# www/nostromo/files/nostromo.sh.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: nostromo
# REQUIRE: DAEMON
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable nostromo:
#
# nostromo_enable="YES"
#

. /etc/rc.subr

name=nostromo
rcvar=nostromo_enable

command=/usr/local/sbin/nhttpd

nostromo_enable=${nostromo_enable:-"NO"}
nostromo_flags=${nostromo_flags:-"-c /usr/local/nostromo/conf/nhttpd.conf"}

load_rc_config $name
run_rc_command "$1"
