# ftp/tftp-hpa/files/openrc-tftpd.in
#!/sbin/openrc-run

name="tftpd"
pidfile=${tftpd_pidfile:-"/var/run/tftpd.pid"}
command="%%PREFIX%%/libexec/in.tftpd"
command_args="-P $pidfile -l"
: ${tftpd_flags:="-s %%PREFIX%%/tftp"}

depend() {
	use net
	keyword -shutdown
}
