# sysutils/appjail/pkg-plist
bin/appjail
bin/appjail-dns
etc/rc.d/appjail
etc/rc.d/appjail-natnet
libexec/appjail/getservbyname/getservbyname
libexec/appjail/ipcheck/ipcheck
libexec/appjail/network/network
%%DATADIR%%/cmd/cmd
%%DATADIR%%/cmd/cpuset
%%DATADIR%%/cmd/disable
%%DATADIR%%/cmd/enable
%%DATADIR%%/cmd/enabled
%%DATADIR%%/cmd/expose
%%DATADIR%%/cmd/fetch
%%DATADIR%%/cmd/fstab
%%DATADIR%%/cmd/help
%%DATADIR%%/cmd/jail
%%DATADIR%%/cmd/limits
%%DATADIR%%/cmd/login
%%DATADIR%%/cmd/logs
%%DATADIR%%/cmd/makejail
%%DATADIR%%/cmd/nat
%%DATADIR%%/cmd/network
%%DATADIR%%/cmd/pkg
%%DATADIR%%/cmd/quick
%%DATADIR%%/cmd/restart
%%DATADIR%%/cmd/rstop
%%DATADIR%%/cmd/run
%%DATADIR%%/cmd/service
%%DATADIR%%/cmd/start
%%DATADIR%%/cmd/startup
%%DATADIR%%/cmd/stop
%%DATADIR%%/cmd/sysrc
%%DATADIR%%/cmd/update
%%DATADIR%%/cmd/upgrade
%%DATADIR%%/cmd/usage
%%DATADIR%%/cmd/version
%%DATADIR%%/cmd/zfs
%%DATADIR%%/files/config.conf
%%DATADIR%%/files/default.conf
%%DATADIR%%/files/default_template.conf
%%DATADIR%%/files/exclude-thinjail.files
%%DATADIR%%/files/include-thinjail.files
%%DATADIR%%/files/rctl-actions
%%DATADIR%%/files/rctl-resources
%%DATADIR%%/files/rctl-subjects
%%DATADIR%%/lib/atexit
%%DATADIR%%/lib/check_func
%%DATADIR%%/lib/colors
%%DATADIR%%/lib/copy
%%DATADIR%%/lib/files
%%DATADIR%%/lib/jail
%%DATADIR%%/lib/jail_types
%%DATADIR%%/lib/kern_modules
%%DATADIR%%/lib/keys
%%DATADIR%%/lib/load
%%DATADIR%%/lib/log
%%DATADIR%%/lib/mksum
%%DATADIR%%/lib/mount
%%DATADIR%%/lib/network
%%DATADIR%%/lib/random
%%DATADIR%%/lib/replace
%%DATADIR%%/lib/strlen
%%DATADIR%%/lib/sysexits
%%DATADIR%%/lib/tempfile
%%DATADIR%%/lib/version
%%DATADIR%%/lib/whitespaces
%%DATADIR%%/lib/zfs
%%DATADIR%%/makejail/cmd/all/ADD
%%DATADIR%%/makejail/cmd/all/ARG
%%DATADIR%%/makejail/cmd/all/CLEAR
%%DATADIR%%/makejail/cmd/all/CMD
%%DATADIR%%/makejail/cmd/all/COPY
%%DATADIR%%/makejail/cmd/all/ENTRYPOINT
%%DATADIR%%/makejail/cmd/all/ENV
%%DATADIR%%/makejail/cmd/all/MOUNT
%%DATADIR%%/makejail/cmd/all/PKG
%%DATADIR%%/makejail/cmd/all/RAW
%%DATADIR%%/makejail/cmd/all/REPLACE
%%DATADIR%%/makejail/cmd/all/RUN
%%DATADIR%%/makejail/cmd/all/SERVICE
%%DATADIR%%/makejail/cmd/all/SET
%%DATADIR%%/makejail/cmd/all/SYSRC
%%DATADIR%%/makejail/cmd/all/UMOUNT
%%DATADIR%%/makejail/cmd/all/USER
%%DATADIR%%/makejail/cmd/all/VAR
%%DATADIR%%/makejail/cmd/all/WORKDIR
%%DATADIR%%/makejail/cmd/build/LOGIN
%%DATADIR%%/makejail/cmd/build/OPTION
%%DATADIR%%/makejail/cmd/build/RESTART
%%DATADIR%%/makejail/cmd/build/START
%%DATADIR%%/makejail/cmd/build/STOP
%%DATADIR%%/makejail/write/all/ADD
%%DATADIR%%/makejail/write/all/ARG
%%DATADIR%%/makejail/write/all/CMD
%%DATADIR%%/makejail/write/all/COPY
%%DATADIR%%/makejail/write/all/MOUNT
%%DATADIR%%/makejail/write/all/PKG
%%DATADIR%%/makejail/write/all/RAW
%%DATADIR%%/makejail/write/all/REPLACE
%%DATADIR%%/makejail/write/all/RUN
%%DATADIR%%/makejail/write/all/SERVICE
%%DATADIR%%/makejail/write/all/SET
%%DATADIR%%/makejail/write/all/SYSRC
%%DATADIR%%/makejail/write/all/UMOUNT
%%DATADIR%%/makejail/write/all/VAR
%%DATADIR%%/makejail/write/all/WORKDIR
%%DATADIR%%/makejail/write/build/LOGIN
%%DATADIR%%/makejail/write/build/OPTION
%%DATADIR%%/makejail/write/build/RESTART
%%DATADIR%%/makejail/write/build/START
%%DATADIR%%/makejail/write/build/STOP
%%DATADIR%%/scripts/ascii2oct.sh
%%DATADIR%%/scripts/dns.sh
%%DATADIR%%/scripts/get_depends.sh
%%DATADIR%%/scripts/run_init.sh
%%DATADIR%%/scripts/run_jail.sh
bin/appjail-config
bin/appjail-config-user
bin/appjail-user
etc/rc.d/appjail-health
libexec/appjail/appjail-config/appjail-config
libexec/appjail/appjail-config/tok
%%DATADIR%%/cmd/healthcheck
%%DATADIR%%/cmd/status
%%DATADIR%%/makejail/cmd/all/DESTROY
%%DATADIR%%/makejail/cmd/all/EXEC
%%DATADIR%%/makejail/write/all/DESTROY
%%DATADIR%%/makejail/write/all/EXEC
%%DATADIR%%/scripts/ajconf-user.sh
%%DATADIR%%/scripts/ajconf.sh
%%DATADIR%%/scripts/ajuser.sh
%%DATADIR%%/scripts/runas.sh
@sample %%ETCDIR%%/appjail.conf.sample
%%DATADIR%%/lib/repeat
%%DATADIR%%/cmd/image
%%DATADIR%%/makejail/cmd/build/FROM
%%DATADIR%%/makejail/write/build/FROM
%%DATADIR%%/scripts/kill_child.sh
etc/rc.d/appjail-dns
%%DATADIR%%/files/dnsmasq.conf
%%DATADIR%%/makejail/cmd/all/GLOBAL
%%DATADIR%%/makejail/write/all/GLOBAL
%%DATADIR%%/scripts/ajdns.sh
%%DATADIR%%/scripts/dnsmasq-hook.sh
libexec/appjail/find-number-from-start/find-number-from-start
libexec/appjail/find-smallest-missing-number/find-smallest-missing-number
%%DATADIR%%/cmd/apply
%%DATADIR%%/cmd/checkOld
%%DATADIR%%/cmd/deleteOld
%%DATADIR%%/cmd/devfs
%%DATADIR%%/cmd/etcupdate
%%DATADIR%%/lib/devfs
%%DATADIR%%/lib/select
%%DATADIR%%/makejail/cmd/all/DEVICE
%%DATADIR%%/makejail/write/all/DEVICE
%%DATADIR%%/scripts/load-devfs-rules.sh
%%DATADIR%%/scripts/super-mv.sh
%%DATADIR%%/scripts/kill_tree.sh
%%DATADIR%%/cmd/volume
%%DATADIR%%/makejail/cmd/all/VOLUME
%%DATADIR%%/makejail/write/all/VOLUME
%%DATADIR%%/scripts/escape-env-val.sh
%%EXAMPLESDIR%%/appjail.conf
%%EXAMPLESDIR%%/templates/bridge.conf
%%EXAMPLESDIR%%/templates/linux.conf
%%EXAMPLESDIR%%/templates/vnet.conf
%%EXAMPLESDIR%%/templates/xorg.conf
