# sysutils/freqsdwn/files/000.freqshpk.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: freqshpk
# REQUIRE: DAEMON
# KEYWORD: nostart shutdown
#
# Define these freqshpk_* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#       /etc/rc.conf.d/freqshpk
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#
freqshpk_enable=${freqshpk_enable-"NO"}

. /etc/rc.subr

name=freqshpk
rcvar=freqshpk_enable
command=%%PREFIX%%/sbin/freqshpk

load_rc_config ${name}

start_cmd=:
required_files=%%PREFIX%%/etc/freqship/UPSFILE
stop_cmd=${command}

run_rc_command "$1"
