# net-mgmt/netmagis-common/files/topographd.in
#! /bin/sh
#
# $FreeBSD$
#

# PROVIDE: topographd
# REQUIRE: LOGIN
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable topographd:
#
# topographd_enable="YES"
#
# # Optionnal variables:
# topographd_flags=""
# topographd_user="rancid"
#

. /etc/rc.subr

name=topographd
rcvar=topographd_enable

load_rc_config $name

: ${topographd_enable:=NO}
: ${topographd_user:=rancid}

command_interpreter=%%TCLSH%%

command=%%PREFIX%%/sbin/topographd

run_rc_command $1
