# www/web2ldap/pkg-plist
%%DATADIR%%/AUTHORS
%%DATADIR%%/CHANGES
%%DATADIR%%/INSTALL
%%DATADIR%%/README
%%DATADIR%%/cgi-bin/web2ldap.cgi
%%DATADIR%%/etc/dumpasn1.cfg
%%DATADIR%%/etc/httpd/sample-mod_fastcgi.conf
%%DATADIR%%/etc/httpd/sample-mod_fcgid.conf
%%DATADIR%%/etc/httpd/sample-mod_scgi.conf
%%DATADIR%%/etc/init.d/web2ldap
%%DATADIR%%/etc/web2ldap/ldapoidreg.py
%%DATADIR%%/etc/web2ldap/ldapoidreg.pyc
%%DATADIR%%/etc/web2ldap/localschema.ldif
%%DATADIR%%/etc/web2ldap/ssl/crt/stroeder.com_Server_CA_no._2008-07.pem
%%DATADIR%%/etc/web2ldap/templates/add_group.ldif
%%DATADIR%%/etc/web2ldap/templates/add_msorganization.ldif
%%DATADIR%%/etc/web2ldap/templates/add_msperson.ldif
%%DATADIR%%/etc/web2ldap/templates/add_olcBdbConfig_AddressBook.ldif
%%DATADIR%%/etc/web2ldap/templates/add_olcHdbConfig_UnixUsers.ldif
%%DATADIR%%/etc/web2ldap/templates/add_olcSchemaConfig.ldif
%%DATADIR%%/etc/web2ldap/templates/add_organization.ldif
%%DATADIR%%/etc/web2ldap/templates/add_orgunit.ldif
%%DATADIR%%/etc/web2ldap/templates/add_person.ldif
%%DATADIR%%/etc/web2ldap/templates/add_user.ldif
%%DATADIR%%/etc/web2ldap/templates/inputform_account.html
%%DATADIR%%/etc/web2ldap/templates/inputform_eduPerson.html
%%DATADIR%%/etc/web2ldap/templates/inputform_inetOrgPerson.html
%%DATADIR%%/etc/web2ldap/templates/inputform_msOrganization.html
%%DATADIR%%/etc/web2ldap/templates/inputform_msPerson.html
%%DATADIR%%/etc/web2ldap/templates/inputform_olcBdbConfig.html
%%DATADIR%%/etc/web2ldap/templates/inputform_olcHdbConfig.html
%%DATADIR%%/etc/web2ldap/templates/inputform_organization.html
%%DATADIR%%/etc/web2ldap/templates/inputform_posixAccount.html
%%DATADIR%%/etc/web2ldap/templates/login.html
%%DATADIR%%/etc/web2ldap/templates/print_organization.html
%%DATADIR%%/etc/web2ldap/templates/print_organizationalUnit.html
%%DATADIR%%/etc/web2ldap/templates/print_person.html
%%DATADIR%%/etc/web2ldap/templates/read_account.html
%%DATADIR%%/etc/web2ldap/templates/read_bankArrangement.html
%%DATADIR%%/etc/web2ldap/templates/read_eduPerson.html
%%DATADIR%%/etc/web2ldap/templates/read_inetOrgPerson.html
%%DATADIR%%/etc/web2ldap/templates/read_inetOrgPerson.html.de
%%DATADIR%%/etc/web2ldap/templates/read_msOrganization.html
%%DATADIR%%/etc/web2ldap/templates/read_msPerson.html
%%DATADIR%%/etc/web2ldap/templates/read_musician.html
%%DATADIR%%/etc/web2ldap/templates/read_musician.html.de
%%DATADIR%%/etc/web2ldap/templates/read_operationalattrs.html
%%DATADIR%%/etc/web2ldap/templates/read_operationalattrs.html.de
%%DATADIR%%/etc/web2ldap/templates/read_organization.html
%%DATADIR%%/etc/web2ldap/templates/read_posixAccount.html
%%DATADIR%%/etc/web2ldap/templates/searchform_Base.html
%%DATADIR%%/etc/web2ldap/templates/searchform_Base.html.de
%%DATADIR%%/etc/web2ldap/templates/searchform_orga.html
%%DATADIR%%/etc/web2ldap/templates/searchform_persons.html
%%DATADIR%%/etc/web2ldap/templates/searchform_users.html
%%DATADIR%%/etc/web2ldap/templates/vcard_organization.txt
%%DATADIR%%/etc/web2ldap/templates/vcard_person.txt
%%DATADIR%%/etc/web2ldap/web2ldapcnf/__init__.py
%%DATADIR%%/etc/web2ldap/web2ldapcnf/__init__.pyc
%%DATADIR%%/etc/web2ldap/web2ldapcnf/countries.py
%%DATADIR%%/etc/web2ldap/web2ldapcnf/countries.pyc
%%DATADIR%%/etc/web2ldap/web2ldapcnf/fastcgi.py
%%DATADIR%%/etc/web2ldap/web2ldapcnf/fastcgi.pyc
%%DATADIR%%/etc/web2ldap/web2ldapcnf/hosts.py
%%DATADIR%%/etc/web2ldap/web2ldapcnf/hosts.pyc
%%DATADIR%%/etc/web2ldap/web2ldapcnf/misc.py
%%DATADIR%%/etc/web2ldap/web2ldapcnf/misc.pyc
%%DATADIR%%/etc/web2ldap/web2ldapcnf/monitor.py
%%DATADIR%%/etc/web2ldap/web2ldapcnf/monitor.pyc
%%DATADIR%%/etc/web2ldap/web2ldapcnf/plugins/__init__.py
%%DATADIR%%/etc/web2ldap/web2ldapcnf/plugins/__init__.pyc
%%DATADIR%%/etc/web2ldap/web2ldapcnf/scgi.py
%%DATADIR%%/etc/web2ldap/web2ldapcnf/scgi.pyc
%%DATADIR%%/etc/web2ldap/web2ldapcnf/standalone.py
%%DATADIR%%/etc/web2ldap/web2ldapcnf/standalone.pyc
%%DATADIR%%/fcgi/web2ldap.py
%%DATADIR%%/fcgi/web2ldap.pyc
%%DATADIR%%/pylib/ipadr.py
%%DATADIR%%/pylib/ipadr.pyc
%%DATADIR%%/pylib/ldapdns.py
%%DATADIR%%/pylib/ldapdns.pyc
%%DATADIR%%/pylib/ldapsession.py
%%DATADIR%%/pylib/ldapsession.pyc
%%DATADIR%%/pylib/ldaputil/__init__.py
%%DATADIR%%/pylib/ldaputil/__init__.pyc
%%DATADIR%%/pylib/ldaputil/base.py
%%DATADIR%%/pylib/ldaputil/base.pyc
%%DATADIR%%/pylib/ldaputil/modlist2.py
%%DATADIR%%/pylib/ldaputil/modlist2.pyc
%%DATADIR%%/pylib/ldaputil/passwd.py
%%DATADIR%%/pylib/ldaputil/passwd.pyc
%%DATADIR%%/pylib/ldaputil/schema.py
%%DATADIR%%/pylib/ldaputil/schema.pyc
%%DATADIR%%/pylib/msHTTPHandler.py
%%DATADIR%%/pylib/msHTTPHandler.pyc
%%DATADIR%%/pylib/msHTTPServer.py
%%DATADIR%%/pylib/msHTTPServer.pyc
%%DATADIR%%/pylib/msbase.py
%%DATADIR%%/pylib/msbase.pyc
%%DATADIR%%/pylib/msgzip.py
%%DATADIR%%/pylib/msgzip.pyc
%%DATADIR%%/pylib/mspki/__init__.py
%%DATADIR%%/pylib/mspki/__init__.pyc
%%DATADIR%%/pylib/mspki/asn1helper.py
%%DATADIR%%/pylib/mspki/asn1helper.pyc
%%DATADIR%%/pylib/mspki/asn1types.py
%%DATADIR%%/pylib/mspki/asn1types.pyc
%%DATADIR%%/pylib/mspki/nsext.py
%%DATADIR%%/pylib/mspki/nsext.pyc
%%DATADIR%%/pylib/mspki/pkix.py
%%DATADIR%%/pylib/mspki/pkix.pyc
%%DATADIR%%/pylib/mspki/utctime.py
%%DATADIR%%/pylib/mspki/utctime.pyc
%%DATADIR%%/pylib/mspki/util.py
%%DATADIR%%/pylib/mspki/util.pyc
%%DATADIR%%/pylib/mspki/vendorext.py
%%DATADIR%%/pylib/mspki/vendorext.pyc
%%DATADIR%%/pylib/mspki/x500.py
%%DATADIR%%/pylib/mspki/x500.pyc
%%DATADIR%%/pylib/mspki/x509.py
%%DATADIR%%/pylib/mspki/x509.pyc
%%DATADIR%%/pylib/mspki/x509v3.py
%%DATADIR%%/pylib/mspki/x509v3.pyc
%%DATADIR%%/pylib/mssignals.py
%%DATADIR%%/pylib/mssignals.pyc
%%DATADIR%%/pylib/pisces/LICENSE
%%DATADIR%%/pylib/pisces/PATCHES
%%DATADIR%%/pylib/pisces/__init__.py
%%DATADIR%%/pylib/pisces/__init__.pyc
%%DATADIR%%/pylib/pisces/asn1.py
%%DATADIR%%/pylib/pisces/asn1.pyc
%%DATADIR%%/pylib/t61_7.py
%%DATADIR%%/pylib/t61_7.pyc
%%DATADIR%%/pylib/t61_8.py
%%DATADIR%%/pylib/t61_8.pyc
%%DATADIR%%/pylib/utctime.py
%%DATADIR%%/pylib/utctime.pyc
%%DATADIR%%/pylib/w2lapp/__init__.py
%%DATADIR%%/pylib/w2lapp/__init__.pyc
%%DATADIR%%/pylib/w2lapp/add.py
%%DATADIR%%/pylib/w2lapp/add.pyc
%%DATADIR%%/pylib/w2lapp/addmodifyform.py
%%DATADIR%%/pylib/w2lapp/addmodifyform.pyc
%%DATADIR%%/pylib/w2lapp/cnf.py
%%DATADIR%%/pylib/w2lapp/cnf.pyc
%%DATADIR%%/pylib/w2lapp/connect.py
%%DATADIR%%/pylib/w2lapp/connect.pyc
%%DATADIR%%/pylib/w2lapp/conninfo.py
%%DATADIR%%/pylib/w2lapp/conninfo.pyc
%%DATADIR%%/pylib/w2lapp/core.py
%%DATADIR%%/pylib/w2lapp/core.pyc
%%DATADIR%%/pylib/w2lapp/delete.py
%%DATADIR%%/pylib/w2lapp/delete.pyc
%%DATADIR%%/pylib/w2lapp/form.py
%%DATADIR%%/pylib/w2lapp/form.pyc
%%DATADIR%%/pylib/w2lapp/groupadm.py
%%DATADIR%%/pylib/w2lapp/groupadm.pyc
%%DATADIR%%/pylib/w2lapp/gui.py
%%DATADIR%%/pylib/w2lapp/gui.pyc
%%DATADIR%%/pylib/w2lapp/handler.py
%%DATADIR%%/pylib/w2lapp/handler.pyc
%%DATADIR%%/pylib/w2lapp/locate.py
%%DATADIR%%/pylib/w2lapp/locate.pyc
%%DATADIR%%/pylib/w2lapp/login.py
%%DATADIR%%/pylib/w2lapp/login.pyc
%%DATADIR%%/pylib/w2lapp/modify.py
%%DATADIR%%/pylib/w2lapp/modify.pyc
%%DATADIR%%/pylib/w2lapp/monitor.py
%%DATADIR%%/pylib/w2lapp/monitor.pyc
%%DATADIR%%/pylib/w2lapp/passwd.py
%%DATADIR%%/pylib/w2lapp/passwd.pyc
%%DATADIR%%/pylib/w2lapp/read.py
%%DATADIR%%/pylib/w2lapp/read.pyc
%%DATADIR%%/pylib/w2lapp/referral.py
%%DATADIR%%/pylib/w2lapp/referral.pyc
%%DATADIR%%/pylib/w2lapp/rename.py
%%DATADIR%%/pylib/w2lapp/rename.pyc
%%DATADIR%%/pylib/w2lapp/schema/__init__.py
%%DATADIR%%/pylib/w2lapp/schema/__init__.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/__init__.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/__init__.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/activedirectory.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/activedirectory.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/asn1objects.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/asn1objects.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/dhcp.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/dhcp.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/dirx.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/dirx.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/edirectory.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/edirectory.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/eduperson.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/eduperson.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/entrust.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/entrust.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/exchange.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/exchange.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/ibmds.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/ibmds.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/krb5.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/krb5.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/lotusdomino.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/lotusdomino.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/msperson.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/msperson.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/mssfu30.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/mssfu30.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/nis.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/nis.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/openldap.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/openldap.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/pgpkeysrv.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/pgpkeysrv.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/pilotperson.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/pilotperson.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/ppolicy.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/ppolicy.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/quirks.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/quirks.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/samba.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/samba.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/schac.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/schac.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/subentries.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/subentries.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/vpim.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/vpim.pyc
%%DATADIR%%/pylib/w2lapp/schema/plugins/x500dsa.py
%%DATADIR%%/pylib/w2lapp/schema/plugins/x500dsa.pyc
%%DATADIR%%/pylib/w2lapp/schema/syntaxes.py
%%DATADIR%%/pylib/w2lapp/schema/syntaxes.pyc
%%DATADIR%%/pylib/w2lapp/schema/viewer.py
%%DATADIR%%/pylib/w2lapp/schema/viewer.pyc
%%DATADIR%%/pylib/w2lapp/search.py
%%DATADIR%%/pylib/w2lapp/search.pyc
%%DATADIR%%/pylib/w2lapp/searchform.py
%%DATADIR%%/pylib/w2lapp/searchform.pyc
%%DATADIR%%/pylib/w2lapp/session.py
%%DATADIR%%/pylib/w2lapp/session.pyc
%%DATADIR%%/pylib/w2lapp/srvrr.py
%%DATADIR%%/pylib/w2lapp/srvrr.pyc
%%DATADIR%%/pylib/w2lapp/viewer.py
%%DATADIR%%/pylib/w2lapp/viewer.pyc
%%DATADIR%%/sbin/checkinst.py
%%DATADIR%%/sbin/checkinst.pyc
%%DATADIR%%/sbin/compile.py
%%DATADIR%%/sbin/compile.pyc
%%DATADIR%%/sbin/web2ldap.py
%%DATADIR%%/sbin/web2ldap.pyc
%%DATADIR%%/scgi/web2ldap.py
%%DATADIR%%/scgi/web2ldap.pyc
%%PORTDOCS%%%%DOCSDIR%%/htdocs/.htaccess
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes-0.10.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes-0.11.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes-0.12.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes-0.13.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes-0.14.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes-0.15.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes-0.16.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes-0.7.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes-0.8.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes-0.9.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes-1.0.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes-1.1.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes-ancient.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/changes.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/compability.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/css/web2ldap/default.css
%%PORTDOCS%%%%DOCSDIR%%/htdocs/css/web2ldap/fixed.css
%%PORTDOCS%%%%DOCSDIR%%/htdocs/css/web2ldap/print.css
%%PORTDOCS%%%%DOCSDIR%%/htdocs/demo.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/docs.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/download/.htaccess
%%PORTDOCS%%%%DOCSDIR%%/htdocs/download/redirect-to-download.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/download.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/faq.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/features.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/feedback.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/files.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/img/ldif_template_select_list.jpg
%%PORTDOCS%%%%DOCSDIR%%/htdocs/install.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/news.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/related.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/roadmap.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/robots.txt
%%PORTDOCS%%%%DOCSDIR%%/htdocs/security.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/sidebars/mozilla-add-sidebar.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/sidebars/mozilla-sidebar.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/ssi/changes_navigation.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/ssi/footer.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/ssi/head.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/ssi/navigation.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/support.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/usability.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/usage.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/web2ldap.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/web2ldapcnf.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/web2ldapcnf_fastcgi.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/web2ldapcnf_hosts.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/web2ldapcnf_misc.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/web2ldapcnf_monitor.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/web2ldapcnf_scgi.html
%%PORTDOCS%%%%DOCSDIR%%/htdocs/web2ldapcnf_standalone.html
@comment directories in revers order
%%PORTDOCS%%@dirrm %%DOCSDIR%%/htdocs/ssi
%%PORTDOCS%%@dirrm %%DOCSDIR%%/htdocs/sidebars
%%PORTDOCS%%@dirrm %%DOCSDIR%%/htdocs/img
%%PORTDOCS%%@dirrm %%DOCSDIR%%/htdocs/download
%%PORTDOCS%%@dirrm %%DOCSDIR%%/htdocs/css/web2ldap
%%PORTDOCS%%@dirrm %%DOCSDIR%%/htdocs/css
%%PORTDOCS%%@dirrm %%DOCSDIR%%/htdocs
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrmtry /var/log/web2ldap
@dirrm %%DATADIR%%/var/run
@dirrm %%DATADIR%%/var/log
@dirrm %%DATADIR%%/var
@dirrm %%DATADIR%%/scgi
@dirrm %%DATADIR%%/sbin
@dirrm %%DATADIR%%/pylib/w2lapp/schema/plugins
@dirrm %%DATADIR%%/pylib/w2lapp/schema
@dirrm %%DATADIR%%/pylib/w2lapp
@dirrm %%DATADIR%%/pylib/pisces
@dirrm %%DATADIR%%/pylib/mspki
@dirrm %%DATADIR%%/pylib/ldaputil
@dirrm %%DATADIR%%/pylib
@dirrm %%DATADIR%%/fcgi
@dirrm %%DATADIR%%/etc/web2ldap/web2ldapcnf/plugins
@dirrm %%DATADIR%%/etc/web2ldap/web2ldapcnf
@dirrm %%DATADIR%%/etc/web2ldap/templates
@dirrm %%DATADIR%%/etc/web2ldap/ssl/crt
@dirrm %%DATADIR%%/etc/web2ldap/ssl
@dirrm %%DATADIR%%/etc/web2ldap
@dirrm %%DATADIR%%/etc/init.d
@dirrm %%DATADIR%%/etc/httpd
@dirrm %%DATADIR%%/etc
@dirrm %%DATADIR%%/cgi-bin
@dirrm %%DATADIR%%
