# www/xshttpd-devel/files/openrc-xshttpd.in
#!/sbin/openrc-run

name="xshttpd"
description=""

pidfile="/var/run/xshttpd.pid"
command="%%PREFIX%%/bin/httpd"

depend(){
	keyword -shutdown
}
