# databases/pgpool-II-30/files/pgpool.in
#!/bin/sh
#
# $FreeBSD$
#        
# PROVIDE: pgpool
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# pgpool_enable (bool):		Set to NO by default.
#				Set it to YES to enable pgpool.
# pgpool_conf (file):		Set location of your config.
#				Default is "%%PREFIX%%/etc/pgpool.conf"
# pgpool_flags (string):	Default is "-f $pgpool_conf"
# pgpool_looger (bool):		Use logger for logging.
# pgpool_looger_pri (string):	logger priority (default "local0.info").
# pgpool_user (uid):		User ID to run as (default nobody)
#

. /etc/rc.subr

name=pgpool
rcvar=${name}_enable

load_rc_config $name

: ${pgpool_enable="NO"}
: ${pgpool_conf="%%PREFIX%%/etc/pgpool.conf"}
: ${pgpool_flags="-f ${pgpool_conf}"}
: ${pgpool_user="nobody"}
: ${pgpool_logger="YES"}
: ${pgpool_logger_pri="local0.info"}

command="%%PREFIX%%/bin/pgpool"
command_args="${pgpool_flags}"
start_cmd="pgpool_start"
stop_cmd="${command} ${command_args} stop"

required_files="${pgpool_conf}"

pgpool_start()
{
    if checkyesno pgpool_logger; then
	su -m $pgpool_user -c "sh -c \"$command $command_args -n 2>&1 | logger -t $name -p $pgpool_logger_pri &\""
    else
	unset start_cmd
	run_rc_command start
    fi
}

run_rc_command "$1"
