# audio/ezstream/files/openrc-ezstream.in
#!/sbin/openrc-run

name="ezstream"

: ${ezstream_flags="-c %%PREFIX%%/etc/${name}/${name}.xml"}

command="%%PREFIX%%/bin/${name}"
command_args=${ezstream_flags}

depend() {
	keyword -shutdown
}
