# www/cblog/files/cblog.in
#!/bin/sh

# $FreeBSD$
#
# cblog.fcgi startup script
#
# PROVIDE: cblog
# REQUIRE: LOGIN
# KEYWORK: shutdown

# Add the following to /etc/rc.conf[.local] to enable this service
#
# cblog_enable=YES
#
# You can fine tune others variables too:
# cblog_fib="NONE"
# cblog_socket="unix:/var/run/cblog/cblog.sock"
# syntax can be : 
# unix:/patch/to/socket
# tcp:IP:PORT (ipv6 and ipv4 supported)
# Use cblog_user to run cblog as user

cblog_setfib() {
	sysctl net.fibs >/dev/null 2>&1 || return 0

	case "$cblog_fib" in
	[Nn][Oo][Nn][Ee])
		;;
	*)
		command="setfib -F ${cblog_fib} ${command}"
		;;
	esac
}

cblog_precmd() {
	cblog_setfib

	install -d -o $cblog_user -m 755 /var/run/cblog
}

. /etc/rc.subr

name="cblog"
rcvar=`set_rcvar`

command="%%PREFIX%%/libexec/cblog.fcgi"
start_precmd="cblog_precmd"

load_rc_config $name
cblog_enable=${cblog_enable:-"NO"}
cblog_fib=${cblog_fib:-"NONE"}
cblog_user=${cblog_user:-"root"}
cblog_socket=${cblog_socker:-"unix:/var/run/cblog/cblog.sock"}

command_args=${cblog_socket}

run_rc_command "$1"
