# dns/opendd/files/opendd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: opendd
# REQUIRE: NETWORKING
#
# Add the following line to /etc/rc.conf[.local] to enable opendd
#
# opendd_enable (bool):		Set to "NO" by default.
#                               Set it to "YES" to enable opendd.

. /etc/rc.subr

name="opendd"
rcvar=${name}_enable

load_rc_config $name

: ${opendd_enable="NO"}
: ${opendd_pidfile="/var/run/opendd.pid"}

command="%%PREFIX%%/sbin/opendd"
required_files="%%ETCDIR%%/opendd.conf"
pidfile="/var/run/opendd.pid"
start_precmd="opendd_check"

opendd_check()
{
        if ! egrep -q -i -E "^runasdaemon.*=.*1$" ${required_files}
	then
	    err 1 "opendd script need "runasdaemon=1" on config file"
	fi
}
				    
run_rc_command "$1"
