# net/mars_nwe/pkg-plist
etc/nwserv.conf
etc/nwserv.conf.example
etc/rc.d/nwserv.sh
sbin/nwserv
sbin/ncpserv
sbin/nwbind
sbin/nwclient
sbin/nwconn
