# math/flexiblas/pkg-plist
bin/flexiblas
bin/flexiblas-config
@sample etc/flexiblasrc.sample
%%ATLAS%%etc/flexiblasrc.d/ATLAS.conf
etc/flexiblasrc.d/BlisPThread.conf
etc/flexiblasrc.d/NETLIB.conf
%%OPENBLAS%%etc/flexiblasrc.d/OpenBLASOpenMP.conf
include/flexiblas/blas_gnu.h
include/flexiblas/blas_intel.h
include/flexiblas/cblas.h
include/flexiblas/flexiblas_api.h
include/flexiblas/flexiblas_backend.h
include/flexiblas/flexiblas_config.h
include/flexiblas/flexiblas_dummy_cblas.h
include/flexiblas/flexiblas_dummy_fortran.h
include/flexiblas/flexiblas_fortran_char_len.h
include/flexiblas/flexiblas_fortran_mangle.h
include/flexiblas/flexiblas_hook_structure.h
include/flexiblas/flexiblas_mgmt.h
include/flexiblas/flexiblas_real_calls.h
include/flexiblas/flexiblas_real_cblas_calls.h
include/flexiblas/flexiblas_real_lapack.h
include/flexiblas/flexiblas_real_lapacke.h
include/flexiblas/lapack.h
include/flexiblas/lapack_f77.h
include/flexiblas/lapacke.h
include/flexiblas/lapacke_64.h
include/flexiblas/lapacke_config.h
include/flexiblas/lapacke_mangling.h
include/flexiblas/lapacke_utils.h
%%ATLAS%%lib/flexiblas/libflexiblas_atlas.so
lib/flexiblas/libflexiblas_blispthread.so
lib/flexiblas/libflexiblas_fallback_lapack.so
lib/flexiblas/libflexiblas_hook_dummy.so
lib/flexiblas/libflexiblas_hook_profile.so
lib/flexiblas/libflexiblas_netlib.so
%%OPENBLAS%%lib/flexiblas/libflexiblas_openblasopenmp.so
lib/libflexiblas.so
lib/libflexiblas.so.3
lib/libflexiblas.so.%%VE%%
lib/libflexiblas_api.so
lib/libflexiblas_api.so.3
lib/libflexiblas_api.so.%%VE%%
lib/libflexiblas_mgmt.so
lib/libflexiblas_mgmt.so.3
lib/libflexiblas_mgmt.so.%%VE%%
libdata/pkgconfig/flexiblas.pc
libdata/pkgconfig/flexiblas_api.pc
share/man/man1/flexiblas.1.gz
share/man/man3/flexiblas_avail.3.gz
share/man/man3/flexiblas_current_backend.3.gz
share/man/man3/flexiblas_get_num_threads.3.gz
share/man/man3/flexiblas_get_version.3.gz
share/man/man3/flexiblas_list.3.gz
share/man/man3/flexiblas_list_loaded.3.gz
share/man/man3/flexiblas_load_backend.3.gz
share/man/man3/flexiblas_load_backend_library.3.gz
share/man/man3/flexiblas_print_avail_backends.3.gz
share/man/man3/flexiblas_print_current_backend.3.gz
share/man/man3/flexiblas_print_loaded_backends.3.gz
share/man/man3/flexiblas_set_num_threads.3.gz
share/man/man3/flexiblas_switch.3.gz
share/man/man7/flexiblas-api.7.gz
