# dns/dualserver/files/openrc-dualserver.in
#!/sbin/openrc-run

name="dualserver"
command=%%PREFIX%%/sbin/dualserver

stop_pre()
{
	killall dualserver
}

depend() {
	keyword -shutdown
}
