# irc/irc/files/openrc-ircd.in
#!/sbin/openrc-run

name="ircd"
pidfile="%%PREFIX%%/var/run/ircd.pid"
command="%%PREFIX%%/sbin/ircd"

depend() {
	use net
	provide ircd
	keyword -shutdown
}
