# security/scannedonly/files/scannedonly.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: scannedonly
# REQUIRE: LOGIN clamd
# BEFORE: mail
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable the scannedonly daemon:
#
# scannedonly_clamav_enable="YES"
# scannedonly_clamav_flags="<set as needed>"
#
# See scannedonlyd_clamav(1) for flags
#

. /etc/rc.subr

name=scannedonly_clamav
rcvar=`set_rcvar`

command=%%PREFIX%%/sbin/scannedonlyd_clamav
pidfile=/var/run/scannedonlyd_clamav.pid
command_args="--pidfile ${pidfile}"

# read settings, set default values
load_rc_config "$name"
: ${scannedonly_clamav_enable="NO"}

run_rc_command "$1"
