# x11/pcdm/files/pcdm.in
#!/sbin/openrc-run
# Copyright (c) 2016 Joe Maloney <jmaloney@ixsystems.com>
# Released under the 2-clause BSD license.

name=PCDMd
command="/usr/local/sbin/PCDMd"
pidfile=/var/run/PCDMd.pid
description="TrueOS Display Manager"

start_pre() {
  # Check if boot-loader set option to skip xorg
  if [ "`kenv noxorg 2>/dev/null`" = "YES" ] ; then
    exit 1
  fi

  if [ -e /var/run/nologin ]; then
	exit 1 
  fi	
}

stop_pre() {
	#Make sure PCDMd knows that we intend to stop X
	${command} stop
}

stop_post() {
	#Now stop all PCDM sessions
	killall PCDM-session
}

depend() {
	need urandom
	after bootmisc syscons moused
}
