# audio/musicpd/files/openrc-musicpd.in
#!/sbin/openrc-run

name=musicpd
description="Music Player Daemon"

command=%%PREFIX%%/bin/musicpd
command_args=%%PREFIX%%/etc/musicpd.conf
required_files=%%PREFIX%%/etc/musicpd.conf
pidfile="/var/mpd/.mpd/pid"

start_pre()
{
        if [ ! -d /var/mpd ] ; then
           mkdir -p /var/mpd/music
	       mkdir -p /var/mdp/.mpd/playlists
		   mkdir /var/mdp/.mpd/log
           chown -R mpd:mpd /var/mpd/
		fi
}

depend() {
	provide musicpd
	need mixer avahi_daemon 
	keyword -shutdown
}
