# net/rtsp_proxy/pkg-plist
etc/qts_proxy.conf
etc/rc.d/rtsp_proxy.sh
sbin/rtsp_proxy
