# net/lanbrowsing/files/lisa.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: lisa
# REQUIRE: LOGIN

# Add the following lines to /etc/rc.conf to enable lisa:
#
# lisa_enable="YES"
#

. /etc/rc.subr

name="lisa"
rcvar=`set_rcvar`

command="%%PREFIX%%/bin/lisa"

load_rc_config $name

: ${lisa_enable="NO"}
: ${lisa_conffile="%%PREFIX%%/etc/lisarc"}

command_args="-q -c ${lisa_conffile}"
required_files="${lisa_conffile}"

run_rc_command "$1"
