# net/AquaGatekeeper/pkg-plist
sbin/AquaGatekeeper
etc/rc.d/agk.sh.sample
etc/AquaGatekeeper.cfg.sample
%%PORTDOCS%%%%DOCSDIR%%/AdministratorGuide.pdf
%%PORTDOCS%%%%DOCSDIR%%/AquaGatekeeper.htm
%%PORTDOCS%%%%DOCSDIR%%/CommandLineOptions.htm
%%PORTDOCS%%%%DOCSDIR%%/ConfigSyntax.htm
%%PORTDOCS%%%%DOCSDIR%%/License.rtf
%%PORTDOCS%%%%DOCSDIR%%/LicenseRus.rtf
%%PORTDOCS%%%%DOCSDIR%%/OperationPrinciples.htm
%%PORTDOCS%%%%DOCSDIR%%/RevisionHistory.htm
%%PORTDOCS%%%%DOCSDIR%%/default.css
%%PORTDOCS%%%%DOCSDIR%%/Rus/AdministratorGuideRus.pdf
%%PORTDOCS%%%%DOCSDIR%%/Rus/AquaGatekeeper.htm
%%PORTDOCS%%%%DOCSDIR%%/Rus/CommandLineOptions.htm
%%PORTDOCS%%%%DOCSDIR%%/Rus/ConfigSyntax.htm
%%PORTDOCS%%%%DOCSDIR%%/Rus/OperationPrinciples.htm
%%PORTDOCS%%%%DOCSDIR%%/Rus/RevisionHistory.htm
%%PORTDOCS%%%%DOCSDIR%%/Rus/default.css
%%EXAMPLESDIR%%/Example.cfg
%%PORTDOCS%%@dirrm %%DOCSDIR%%/Rus
%%PORTDOCS%%@dirrm %%DOCSDIR%%
%%PORTDOCS%%@dirrm %%EXAMPLESDIR%%
