# devel/hadoop/files/openrc-namenode.in
#!/sbin/openrc-run
#
# Add the following lines to /etc/rc.conf to enable this service:
#
# namenode_user (str):           Set to %%HDFS_USER%% by default.
#                                Set to a valid user to change default 
#                                namenode user.
# namenode_java_home (str):	 Unset by default.
#                                Switch to another JVM.
#
# namenode_log_dir (str):	 Unset by default.
# namenode_pid_dir (str):	 Unset by default.

name=namenode

: ${namenode_user:=%%HDFS_USER%%}
: ${namenode_group:=%%HADOOP_GROUP%%}

command="%%DATADIR%%/bin/hadoop-daemon.sh"
command_args='--config ${HADOOP_CONF_DIR} start namenode'
command_user=${namenode_user}

start_pre()
{
	if [ -n "${namenode_java_home}" ]
	then
		export JAVA_HOME="${namenode_java_home}"
	fi

	for f in %%ETCDIR%%/envvars.d/*.env
	do
		. ${f}
	done
}

stop_pre()
{
	start_pre
}

stop()
{
	su -m ${namenode_user} -c "${command} --config ${HADOOP_CONF_DIR} stop namenode"
}

depend() {
	keyword -shutdown
}
