# security/kanidm/files/kanidmd.in
#!/bin/sh

# PROVIDE: kanidmd
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add these lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# kanidmd_enable (bool):	Set to NO by default.
#				Set it to YES to enable kanidm_unixd.

. /etc/rc.subr

name=kanidmd
rcvar=kanidmd_enable

load_rc_config $name

: ${kanidmd_enable:="NO"}

pidfile="/var/run/kanidmd.pid"
command=/usr/sbin/daemon
command_args="-u %%USER%% -p /var/run/kanidmd.pid -T kanidmd %%PREFIX%%/libexec/${name}"
procname=%%PREFIX%%/libexec/${name}

run_rc_command "$1"
