# devel/subversion/pkg-plist
%%TOOLS%%bin/fsfs-access-map
bin/svn
%%TOOLS%%bin/svnbench
%%TOOLS%%bin/svn-populate-node-origins-index
bin/svnadmin
%%TOOLS%%bin/svnauthz
%%TOOLS%%bin/svnauthz-validate
%%TOOLS%%bin/svndiff
%%TOOLS%%bin/svndiff3
%%TOOLS%%bin/svndiff4
bin/svndumpfilter
bin/svnfsfs
bin/svnlook
bin/svnmucc
%%TOOLS%%bin/svnraisetreeconflict
bin/svnrdump
%%NO_SVNSERVE_WRAPPER%%bin/svnserve
%%SVNSERVE_WRAPPER%%bin/svnserve.bin
%%SVNSERVE_WRAPPER%%@sample bin/svnserve.sample
bin/svnsync
bin/svnversion
%%TOOLS%%bin/x509-parser
include/subversion-1/mod_authz_svn.h
include/subversion-1/mod_dav_svn.h
include/subversion-1/svn-revision.txt
include/subversion-1/svn_auth.h
include/subversion-1/svn_base64.h
include/subversion-1/svn_cache_config.h
include/subversion-1/svn_checksum.h
include/subversion-1/svn_client.h
include/subversion-1/svn_cmdline.h
include/subversion-1/svn_compat.h
include/subversion-1/svn_config.h
include/subversion-1/svn_ctype.h
include/subversion-1/svn_dav.h
include/subversion-1/svn_delta.h
include/subversion-1/svn_diff.h
include/subversion-1/svn_dirent_uri.h
include/subversion-1/svn_dso.h
include/subversion-1/svn_error.h
include/subversion-1/svn_error_codes.h
include/subversion-1/svn_fs.h
include/subversion-1/svn_hash.h
include/subversion-1/svn_io.h
include/subversion-1/svn_iter.h
include/subversion-1/svn_md5.h
include/subversion-1/svn_mergeinfo.h
include/subversion-1/svn_nls.h
include/subversion-1/svn_opt.h
include/subversion-1/svn_path.h
include/subversion-1/svn_pools.h
include/subversion-1/svn_props.h
include/subversion-1/svn_quoprint.h
include/subversion-1/svn_ra.h
include/subversion-1/svn_ra_svn.h
include/subversion-1/svn_repos.h
include/subversion-1/svn_sorts.h
include/subversion-1/svn_string.h
include/subversion-1/svn_subst.h
include/subversion-1/svn_time.h
include/subversion-1/svn_types.h
include/subversion-1/svn_user.h
include/subversion-1/svn_utf.h
include/subversion-1/svn_version.h
include/subversion-1/svn_wc.h
include/subversion-1/svn_x509.h
include/subversion-1/svn_xml.h
lib/libsvn_client-1.a
%%NO_STATIC%%lib/libsvn_client-1.so
%%NO_STATIC%%lib/libsvn_client-1.so.0
%%NO_STATIC%%lib/libsvn_client-1.so.0.0.0
lib/libsvn_delta-1.a
%%NO_STATIC%%lib/libsvn_delta-1.so
%%NO_STATIC%%lib/libsvn_delta-1.so.0
%%NO_STATIC%%lib/libsvn_delta-1.so.0.0.0
lib/libsvn_diff-1.a
%%NO_STATIC%%lib/libsvn_diff-1.so
%%NO_STATIC%%lib/libsvn_diff-1.so.0
%%NO_STATIC%%lib/libsvn_diff-1.so.0.0.0
lib/libsvn_fs-1.a
%%NO_STATIC%%lib/libsvn_fs-1.so
%%NO_STATIC%%lib/libsvn_fs-1.so.0
%%NO_STATIC%%lib/libsvn_fs-1.so.0.0.0
%%BDB%%lib/libsvn_fs_base-1.a
%%NO_STATIC%%%%BDB%%lib/libsvn_fs_base-1.so
%%NO_STATIC%%%%BDB%%lib/libsvn_fs_base-1.so.0
%%NO_STATIC%%%%BDB%%lib/libsvn_fs_base-1.so.0.0.0
lib/libsvn_fs_fs-1.a
%%NO_STATIC%%lib/libsvn_fs_fs-1.so
%%NO_STATIC%%lib/libsvn_fs_fs-1.so.0
%%NO_STATIC%%lib/libsvn_fs_fs-1.so.0.0.0
lib/libsvn_fs_util-1.a
%%NO_STATIC%%lib/libsvn_fs_util-1.so
%%NO_STATIC%%lib/libsvn_fs_util-1.so.0
%%NO_STATIC%%lib/libsvn_fs_util-1.so.0.0.0
lib/libsvn_fs_x-1.a
%%NO_STATIC%%lib/libsvn_fs_x-1.so
%%NO_STATIC%%lib/libsvn_fs_x-1.so.0
%%NO_STATIC%%lib/libsvn_fs_x-1.so.0.0.0
lib/libsvn_ra-1.a
%%NO_STATIC%%lib/libsvn_ra-1.so
%%NO_STATIC%%lib/libsvn_ra-1.so.0
%%NO_STATIC%%lib/libsvn_ra-1.so.0.0.0
lib/libsvn_ra_local-1.a
%%NO_STATIC%%lib/libsvn_ra_local-1.so
%%NO_STATIC%%lib/libsvn_ra_local-1.so.0
%%NO_STATIC%%lib/libsvn_ra_local-1.so.0.0.0
%%SERF%%lib/libsvn_ra_serf-1.a
%%NO_STATIC%%%%SERF%%lib/libsvn_ra_serf-1.so
%%NO_STATIC%%%%SERF%%lib/libsvn_ra_serf-1.so.0
%%NO_STATIC%%%%SERF%%lib/libsvn_ra_serf-1.so.0.0.0
lib/libsvn_ra_svn-1.a
%%NO_STATIC%%lib/libsvn_ra_svn-1.so
%%NO_STATIC%%lib/libsvn_ra_svn-1.so.0
%%NO_STATIC%%lib/libsvn_ra_svn-1.so.0.0.0
lib/libsvn_repos-1.a
%%NO_STATIC%%lib/libsvn_repos-1.so
%%NO_STATIC%%lib/libsvn_repos-1.so.0
%%NO_STATIC%%lib/libsvn_repos-1.so.0.0.0
lib/libsvn_subr-1.a
%%NO_STATIC%%lib/libsvn_subr-1.so
%%NO_STATIC%%lib/libsvn_subr-1.so.0
%%NO_STATIC%%lib/libsvn_subr-1.so.0.0.0
lib/libsvn_wc-1.a
%%NO_STATIC%%lib/libsvn_wc-1.so
%%NO_STATIC%%lib/libsvn_wc-1.so.0
%%NO_STATIC%%lib/libsvn_wc-1.so.0.0.0
man/man1/svn.1.gz
man/man1/svnadmin.1.gz
man/man1/svndumpfilter.1.gz
man/man1/svnlook.1.gz
man/man1/svnmucc.1.gz
man/man1/svnrdump.1.gz
man/man1/svnsync.1.gz
man/man1/svnversion.1.gz
man/man5/svnserve.conf.5.gz
man/man8/svnserve.8.gz
%%NLS%%share/locale/de/LC_MESSAGES/subversion.mo
%%NLS%%share/locale/es/LC_MESSAGES/subversion.mo
%%NLS%%share/locale/fr/LC_MESSAGES/subversion.mo
%%NLS%%share/locale/it/LC_MESSAGES/subversion.mo
%%NLS%%share/locale/ja/LC_MESSAGES/subversion.mo
%%NLS%%share/locale/ko/LC_MESSAGES/subversion.mo
%%NLS%%share/locale/nb/LC_MESSAGES/subversion.mo
%%NLS%%share/locale/pl/LC_MESSAGES/subversion.mo
%%NLS%%share/locale/pt_BR/LC_MESSAGES/subversion.mo
%%NLS%%share/locale/sv/LC_MESSAGES/subversion.mo
%%NLS%%share/locale/zh_CN/LC_MESSAGES/subversion.mo
%%NLS%%share/locale/zh_TW/LC_MESSAGES/subversion.mo
share/pkgconfig/libsvn_client.pc
share/pkgconfig/libsvn_delta.pc
share/pkgconfig/libsvn_diff.pc
share/pkgconfig/libsvn_fs.pc
share/pkgconfig/libsvn_fs_fs.pc
share/pkgconfig/libsvn_fs_util.pc
share/pkgconfig/libsvn_fs_x.pc
share/pkgconfig/libsvn_ra.pc
share/pkgconfig/libsvn_ra_local.pc
share/pkgconfig/libsvn_ra_serf.pc
share/pkgconfig/libsvn_ra_svn.pc
share/pkgconfig/libsvn_repos.pc
share/pkgconfig/libsvn_subr.pc
share/pkgconfig/libsvn_wc.pc
%%DATADIR%%/backup/hot-backup.py
%%DATADIR%%/bdb/erase-all-text-data.py
%%DATADIR%%/bdb/skel.py
%%DATADIR%%/bdb/svn-bdb-view.py
%%DATADIR%%/bdb/svnfs.py
%%DATADIR%%/bdb/whatis-rep.py
%%DATADIR%%/buildbot/master/README
%%DATADIR%%/buildbot/slaves/README
%%DATADIR%%/buildbot/slaves/bb-openbsd/svnbuild.sh
%%DATADIR%%/buildbot/slaves/bb-openbsd/svncheck-bindings.sh
%%DATADIR%%/buildbot/slaves/bb-openbsd/svncheck.sh
%%DATADIR%%/buildbot/slaves/bb-openbsd/svnclean.sh
%%DATADIR%%/buildbot/slaves/i686-debian-sarge1/mount-ramdrive.c
%%DATADIR%%/buildbot/slaves/i686-debian-sarge1/svnbuild.sh
%%DATADIR%%/buildbot/slaves/i686-debian-sarge1/svncheck.sh
%%DATADIR%%/buildbot/slaves/i686-debian-sarge1/svnclean.sh
%%DATADIR%%/buildbot/slaves/i686-debian-sarge1/svnlog.sh
%%DATADIR%%/buildbot/slaves/i686-debian-sarge1/unmount-ramdrive.c
%%DATADIR%%/buildbot/slaves/svn-sparc-solaris/svnbuild.sh
%%DATADIR%%/buildbot/slaves/svn-sparc-solaris/svncheck.sh
%%DATADIR%%/buildbot/slaves/svn-sparc-solaris/svncleanup.sh
%%DATADIR%%/buildbot/slaves/svn-sparc-solaris/svnenv.sh
%%DATADIR%%/buildbot/slaves/svn-x64-centos/list-svn-deps.sh
%%DATADIR%%/buildbot/slaves/svn-x64-centos/svnbuild.sh
%%DATADIR%%/buildbot/slaves/svn-x64-centos/svncheck-bindings.sh
%%DATADIR%%/buildbot/slaves/svn-x64-centos/svncheck.sh
%%DATADIR%%/buildbot/slaves/svn-x64-centos/svnclean.sh
%%DATADIR%%/buildbot/slaves/svn-x64-centos/svnlog.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx-gnu-shared-daily-ra_serf/svnbuild.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx-gnu-shared-daily-ra_serf/svncheck.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx-gnu-shared-daily-ra_serf/svnclean.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx-gnu-shared-daily-ra_serf/svnlog.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx-gnu-shared/svnbuild.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx-gnu-shared/svncheck.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx-gnu-shared/svnclean.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx-gnu-shared/svnlog.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx/mkramdisk.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx/rmramdisk.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx/setenv.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx/svnbuild-bindings.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx/svnbuild.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx/svncheck-bindings.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx/svncheck.sh
%%DATADIR%%/buildbot/slaves/svn-x64-macosx/svnclean.sh
%%DATADIR%%/buildbot/slaves/ubuntu-x64/svnbuild.sh
%%DATADIR%%/buildbot/slaves/ubuntu-x64/svncheck-bindings.sh
%%DATADIR%%/buildbot/slaves/ubuntu-x64/svncheck.sh
%%DATADIR%%/buildbot/slaves/ubuntu-x64/svnclean.sh
%%DATADIR%%/buildbot/slaves/ubuntu-x64/svnlog.sh
%%DATADIR%%/buildbot/slaves/win32-SharpSvn/svn-config.cmd.template
%%DATADIR%%/buildbot/slaves/win32-SharpSvn/svntest-bindings.cmd
%%DATADIR%%/buildbot/slaves/win32-SharpSvn/svntest-build-bindings.cmd
%%DATADIR%%/buildbot/slaves/win32-SharpSvn/svntest-build.cmd
%%DATADIR%%/buildbot/slaves/win32-SharpSvn/svntest-cleanup.cmd
%%DATADIR%%/buildbot/slaves/win32-SharpSvn/svntest-javahl.cmd
%%DATADIR%%/buildbot/slaves/win32-SharpSvn/svntest-template.cmd
%%DATADIR%%/buildbot/slaves/win32-SharpSvn/svntest-test.cmd
%%DATADIR%%/buildbot/slaves/win32-xp-VS2005/config.bat
%%DATADIR%%/buildbot/slaves/win32-xp-VS2005/do_all.bat
%%DATADIR%%/buildbot/slaves/win32-xp-VS2005/svnbuild.bat
%%DATADIR%%/buildbot/slaves/win32-xp-VS2005/svncheck.bat
%%DATADIR%%/buildbot/slaves/win32-xp-VS2005/svnclean.bat
%%DATADIR%%/buildbot/slaves/win32-xp-VS2005/svnlog.bat
%%DATADIR%%/buildbot/slaves/xp-vc60-ia32/config.bat.tmpl
%%DATADIR%%/buildbot/slaves/xp-vc60-ia32/svnbuild.bat
%%DATADIR%%/buildbot/slaves/xp-vc60-ia32/svncheck.bat
%%DATADIR%%/buildbot/slaves/xp-vc60-ia32/svnclean.bat
%%DATADIR%%/buildbot/slaves/xp-vc60-ia32/svnlog.bat
%%DATADIR%%/client-side/bash_completion
%%DATADIR%%/client-side/bash_completion_test
%%DATADIR%%/client-side/change-svn-wc-format.py
%%DATADIR%%/client-side/detach.py
%%DATADIR%%/client-side/mergeinfo-sanitizer.py
%%DATADIR%%/client-side/server-version.py
%%DATADIR%%/client-side/svn-graph.pl
%%DATADIR%%/client-side/svn-ssl-fingerprints.sh
%%DATADIR%%/client-side/svn-vendor.py
%%DATADIR%%/client-side/svn-viewspec.py
%%DATADIR%%/client-side/wcfind
%%DATADIR%%/dev/analyze-svnlogs.py
%%DATADIR%%/dev/aprerr.txt
%%DATADIR%%/dev/benchmarks/RepoPerf/ClearMemory.cpp
%%DATADIR%%/dev/benchmarks/RepoPerf/TimeWin.cpp
%%DATADIR%%/dev/benchmarks/RepoPerf/copy_repo.py
%%DATADIR%%/dev/benchmarks/RepoPerf/win_repo_bench.py
%%DATADIR%%/dev/benchmarks/large_dirs/create_bigdir.sh
%%DATADIR%%/dev/benchmarks/suite1/benchmark.py
%%DATADIR%%/dev/benchmarks/suite1/cronjob
%%DATADIR%%/dev/benchmarks/suite1/crontab.entry
%%DATADIR%%/dev/benchmarks/suite1/generate_charts
%%DATADIR%%/dev/benchmarks/suite1/run
%%DATADIR%%/dev/benchmarks/suite1/run.bat
%%DATADIR%%/dev/build-svn-deps-win.pl
%%DATADIR%%/dev/check-license.py
%%DATADIR%%/dev/contribulyze.py
%%DATADIR%%/dev/datecheck.py
%%DATADIR%%/dev/find-bad-style.py
%%DATADIR%%/dev/find-unmoved-deprecated.sh
%%DATADIR%%/dev/fsfs-access-map.c
%%DATADIR%%/dev/gdb-py/README
%%DATADIR%%/dev/gdb-py/svndbg/__init__.py
%%DATADIR%%/dev/gdb-py/svndbg/printers.py
%%DATADIR%%/dev/gen-javahl-errors.py
%%DATADIR%%/dev/gen-py-errors.py
%%DATADIR%%/dev/gen_junit_report.py
%%DATADIR%%/dev/gnuify-changelog.pl
%%DATADIR%%/dev/graph-dav-servers.py
%%DATADIR%%/dev/histogram.py
%%DATADIR%%/dev/iz/defect.dem
%%DATADIR%%/dev/iz/ff2csv.command
%%DATADIR%%/dev/iz/ff2csv.py
%%DATADIR%%/dev/iz/find-fix.py
%%DATADIR%%/dev/iz/run-queries.sh
%%DATADIR%%/dev/lock-check.py
%%DATADIR%%/dev/log_revnum_change_asf.py
%%DATADIR%%/dev/merge-graph.py
%%DATADIR%%/dev/mergegraph/__init__.py
%%DATADIR%%/dev/mergegraph/mergegraph.py
%%DATADIR%%/dev/mergegraph/save_as_sh.py
%%DATADIR%%/dev/min-includes.sh
%%DATADIR%%/dev/mklog.py
%%DATADIR%%/dev/mlpatch.py
%%DATADIR%%/dev/normalize-dump.py
%%DATADIR%%/dev/po-merge.py
%%DATADIR%%/dev/prebuild-cleanup.sh
%%DATADIR%%/dev/random-commits.py
%%DATADIR%%/dev/remove-trailing-whitespace.sh
%%DATADIR%%/dev/sbox-ospath.py
%%DATADIR%%/dev/scramble-tree.py
%%DATADIR%%/dev/stress.pl
%%DATADIR%%/dev/svn-dev.el
%%DATADIR%%/dev/svn-dev.vim
%%DATADIR%%/dev/svn-entries.el
%%DATADIR%%/dev/svn-merge-revs.py
%%DATADIR%%/dev/svnqlite3-dump
%%DATADIR%%/dev/svnraisetreeconflict/svnraisetreeconflict.c
%%DATADIR%%/dev/trails.py
%%DATADIR%%/dev/unix-build/Makefile.svn
%%DATADIR%%/dev/unix-build/README
%%DATADIR%%/dev/verify-history.py
%%DATADIR%%/dev/warn-ignored-err.sh
%%DATADIR%%/dev/wc-format.py
%%DATADIR%%/dev/wc-ng/bump-to-19.py
%%DATADIR%%/dev/wc-ng/count-progress.py
%%DATADIR%%/dev/wc-ng/gather-data.sh
%%DATADIR%%/dev/wc-ng/graph-data.py
%%DATADIR%%/dev/wc-ng/populate-pristine.py
%%DATADIR%%/dev/wc-ng/svn-wc-db-tester.c
%%DATADIR%%/dev/which-error.py
%%DATADIR%%/dev/windows-build/Makefile
%%DATADIR%%/dev/windows-build/README
%%DATADIR%%/dev/windows-build/document-version.pl
%%DATADIR%%/dev/x509-parser.c
%%DATADIR%%/diff/diff.c
%%DATADIR%%/diff/diff3.c
%%DATADIR%%/diff/diff4.c
%%DATADIR%%/dist/_gnupg.py
%%DATADIR%%/dist/backport.pl
%%DATADIR%%/dist/backport_accept.dump
%%DATADIR%%/dist/backport_branches.dump
%%DATADIR%%/dist/backport_indented_entry.dump
%%DATADIR%%/dist/backport_multirevisions.dump
%%DATADIR%%/dist/backport_tests.py
%%DATADIR%%/dist/backport_two_approveds.dump
%%DATADIR%%/dist/checksums.py
%%DATADIR%%/dist/collect_sigs.py
%%DATADIR%%/dist/dist.sh
%%DATADIR%%/dist/download-release.sh
%%DATADIR%%/dist/extract-for-examination.sh
%%DATADIR%%/dist/nightly.sh
%%DATADIR%%/dist/nominate.pl
%%DATADIR%%/dist/rat-excludes
%%DATADIR%%/dist/release.py
%%DATADIR%%/dist/templates/download.ezt
%%DATADIR%%/dist/templates/nightly-candidates.ezt
%%DATADIR%%/dist/templates/rc-news.ezt
%%DATADIR%%/dist/templates/rc-release-ann.ezt
%%DATADIR%%/dist/templates/stable-news.ezt
%%DATADIR%%/dist/templates/stable-release-ann.ezt
%%DATADIR%%/dist/test.sh
%%DATADIR%%/examples/ExampleAuthn.java
%%DATADIR%%/examples/ExampleAuthnOld.java
%%DATADIR%%/examples/ExampleAuthnVeryOld.java
%%DATADIR%%/examples/SvnCLBrowse
%%DATADIR%%/examples/blame.py
%%DATADIR%%/examples/check-modified.py
%%DATADIR%%/examples/dumpprops.py
%%DATADIR%%/examples/get-location-segments.py
%%DATADIR%%/examples/getfile.py
%%DATADIR%%/examples/getlocks_test.c
%%DATADIR%%/examples/geturl.py
%%DATADIR%%/examples/headrev.c
%%DATADIR%%/examples/info.rb
%%DATADIR%%/examples/minimal_client.c
%%DATADIR%%/examples/putfile.py
%%DATADIR%%/examples/revplist.py
%%DATADIR%%/examples/svnlog2html.rb
%%DATADIR%%/examples/svnlook.py
%%DATADIR%%/examples/svnlook.rb
%%DATADIR%%/examples/svnput.c
%%DATADIR%%/examples/svnserve-sgid.c
%%DATADIR%%/examples/svnshell.py
%%DATADIR%%/examples/svnshell.rb
%%DATADIR%%/examples/testwrite.c
%%DATADIR%%/examples/walk-config-auth.py
%%DATADIR%%/hook-scripts/commit-access-control.cfg.example
%%DATADIR%%/hook-scripts/commit-access-control.pl
%%DATADIR%%/hook-scripts/commit-email.rb
%%DATADIR%%/hook-scripts/control-chars.py
%%DATADIR%%/hook-scripts/log-police.py
%%DATADIR%%/hook-scripts/mailer/mailer.conf.example
%%DATADIR%%/hook-scripts/mailer/mailer.py
%%DATADIR%%/hook-scripts/mailer/tests/mailer-init.sh
%%DATADIR%%/hook-scripts/mailer/tests/mailer-t1.output
%%DATADIR%%/hook-scripts/mailer/tests/mailer-t1.sh
%%DATADIR%%/hook-scripts/mailer/tests/mailer-tweak.py
%%DATADIR%%/hook-scripts/mailer/tests/mailer.conf
%%DATADIR%%/hook-scripts/persist-ephemeral-txnprops.py
%%DATADIR%%/hook-scripts/svn2feed.py
%%DATADIR%%/hook-scripts/svnperms.conf.example
%%DATADIR%%/hook-scripts/svnperms.py
%%DATADIR%%/hook-scripts/validate-extensions.py
%%DATADIR%%/hook-scripts/validate-files.conf.example
%%DATADIR%%/hook-scripts/validate-files.py
%%DATADIR%%/hook-scripts/verify-po.py
%%DATADIR%%/po/l10n-report.py
%%DATADIR%%/po/po-update.sh
%%DATADIR%%/server-side/fsfs-reshard.py
%%DATADIR%%/server-side/mod_dontdothat/README
%%DATADIR%%/server-side/mod_dontdothat/mod_dontdothat.c
%%DATADIR%%/server-side/svn-backup-dumps.py
%%DATADIR%%/server-side/svn-populate-node-origins-index.c
%%DATADIR%%/server-side/svn_server_log_parse.py
%%DATADIR%%/server-side/svnauthz.c
%%DATADIR%%/server-side/svnpredumpfilter.py
%%DATADIR%%/server-side/svnpubsub/README.txt
%%DATADIR%%/server-side/svnpubsub/commit-hook.py
%%DATADIR%%/server-side/svnpubsub/daemonize.py
%%DATADIR%%/server-side/svnpubsub/irkerbridge.py
%%DATADIR%%/server-side/svnpubsub/rc.d/svnpubsub
%%DATADIR%%/server-side/svnpubsub/rc.d/svnpubsub.debian
%%DATADIR%%/server-side/svnpubsub/rc.d/svnpubsub.freebsd
%%DATADIR%%/server-side/svnpubsub/rc.d/svnpubsub.solaris
%%DATADIR%%/server-side/svnpubsub/rc.d/svnwcsub
%%DATADIR%%/server-side/svnpubsub/rc.d/svnwcsub.debian
%%DATADIR%%/server-side/svnpubsub/rc.d/svnwcsub.freebsd
%%DATADIR%%/server-side/svnpubsub/rc.d/svnwcsub.solaris
%%DATADIR%%/server-side/svnpubsub/revprop-change-hook.py
%%DATADIR%%/server-side/svnpubsub/svnpubsub.tac
%%DATADIR%%/server-side/svnpubsub/svnpubsub/__init__.py
%%DATADIR%%/server-side/svnpubsub/svnpubsub/client.py
%%DATADIR%%/server-side/svnpubsub/svnpubsub/server.py
%%DATADIR%%/server-side/svnpubsub/svnpubsub/util.py
%%DATADIR%%/server-side/svnpubsub/svntweet.py
%%DATADIR%%/server-side/svnpubsub/svnwcsub.conf.example
%%DATADIR%%/server-side/svnpubsub/svnwcsub.py
%%DATADIR%%/server-side/svnpubsub/testserver.py
%%DATADIR%%/server-side/svnpubsub/watcher.py
%%DATADIR%%/server-side/test_svn_server_log_parse.py
%%DATADIR%%/xslt/svnindex.css
%%DATADIR%%/xslt/svnindex.xsl
