# games/fsgs/pkg-plist
etc/rc.d/fsgs.sh
fsgs/Changes.txt
fsgs/License.txt
fsgs/Readme.Linux
fsgs/banners/ad000001.smk
fsgs/banners/banners.ini
fsgs/banners/banners.txt
fsgs/client/readme
fsgs/contrib/fsgsquery.zip
fsgs/contrib/updquery.zip
fsgs/data/D2DV_IX86_102_103.MPQ
fsgs/data/IX86ver1.mpq
fsgs/data/PMACver4.mpq
fsgs/data/adffffff.bik
fsgs/data/adffffff.pcx
fsgs/data/bw4-5.mpq
fsgs/data/bw5-6.mpq
fsgs/data/bw6-7b.mpq
fsgs/data/help.txt
fsgs/data/icons.bni
fsgs/data/icons.ngd
fsgs/data/icons.ngi
fsgs/data/sc4-5.mpq
fsgs/data/sc5-6.mpq
fsgs/data/sc6-7b.mpq
fsgs/data/star107b.mpq
fsgs/data/teams.dat
fsgs/data/tos.txt
fsgs/data/versions.dat
fsgs/fsgs
fsgs/fsgs.cfg
fsgs/fsgs_db.so
fsgs/fsgs_static
fsgs/fsgsprune
fsgs/fsgsprune_static
fsgs/ladder/entry.tpl
fsgs/ladder/footer.tpl
fsgs/ladder/header.tpl
fsgs/logs/df
fsgs/players/df
fsgs/readme.txt
fsgs/reports/report.tpl
@dirrm fsgs/reports
@dirrm fsgs/players
@dirrm fsgs/logs
@dirrm fsgs/ladder
@dirrm fsgs/data
@dirrm fsgs/contrib
@dirrm fsgs/client
@dirrm fsgs/banners
@dirrm fsgs
