# net/xrdp/files/openrc-xrdp.in
#!/sbin/openrc-run

name="xrdp"
command="%%PREFIX%%/sbin/${name}"
command_args="--nodaemon"
pidfile="/var/run/${name}.pid"

start_pre() {
    if [ -f %%PREFIX%%/etc/xrdp/rsakeys.ini ] ; then
        %%PREFIX%%/bin/xrdp-keygen xrdp %%PREFIX%%/etc/xrdp/rsakeys.ini
    fi
}

depend() {
	after net-online
}
