# net/xrdp/files/openrc-xrdp-sesman.in
#!/sbin/openrc-run

name="xrdp-sesman"
command="%%PREFIX%%/sbin/${name}"
command_args="--nodaemon"
pidfile="/var/run/${name}.pid"


depend() {
	need xrdp
}
