# security/strongswan/files/openrc-strongswan.in
#!/sbin/openrc-run

command="%%PREFIX%%/sbin/ipsec"
command_args=start
name="strongswan"

stop_pre() {
	$command stop
}
