# security/strongswan/files/strongswan.in
#!/bin/sh
# Start or stop strongswan
# $FreeBSD$

# PROVIDE: strongswan
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name=strongswan
rcvar=strongswan_enable

load_rc_config $name

: ${strongswan_enable:=NO}

extra_commands="reload statusall"

command="%%PREFIX%%/sbin/ipsec"

start_precmd="strongswan_precmd"
stop_cmd="strongswan_cmd"
status_cmd="strongswan_cmd"
reload_cmd="strongswan_cmd"
statusall_cmd="strongswan_cmd"

strongswan_precmd()
{
	command_args=${rc_arg}
}

strongswan_cmd()
{
	${command} ${rc_arg}
}

run_rc_command "$1"
