# net/knxd/files/knxd.in
#!/bin/sh

# PROVIDE: knxd
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown
#
# Execute the following line to enable knxd:
#
# sysrc knxd_enable="YES":	Set to "YES" to enable knxd
#
# sysrc knxd_config="%%PREFIX%%/etc/knxd.conf"	The config file to use
#

. /etc/rc.subr

name=knxd
rcvar=knxd_enable

load_rc_config $name

: ${knxd_enable:="NO"}
: ${knxd_config:="%%PREFIX%%/etc/knxd.conf"}
: ${knxd_flags="${knxd_config}"}
: ${knxd_svcj_options:="net_basic"}

command="%%PREFIX%%/bin/knxd"
#command_args=""
#stop_postcmd="rm -f $pidfile"

run_rc_command "$1"
