# dns/p5-Zonemaster-Backend/pkg-plist
@dir(zonemaster,zonemaster,750) /var/db/zonemaster
@dir(zonemaster,zonemaster,750) /var/log/zonemaster
@dir(zonemaster,zonemaster,750) /var/run/zonemaster
bin/zmb
bin/zmtest
bin/zonemaster_backend_rpcapi.psgi
bin/zonemaster_backend_testagent
etc/rc.d/zm_rpcapi
etc/rc.d/zm_testagent
@sample etc/zonemaster/backend_config.ini.sample
%%SITE_PERL%%/Zonemaster/Backend.pm
%%SITE_PERL%%/Zonemaster/Backend/Config.pm
%%SITE_PERL%%/Zonemaster/Backend/Config/DCPlugin.pm
%%SITE_PERL%%/Zonemaster/Backend/DB.pm
%%SITE_PERL%%/Zonemaster/Backend/DB/MySQL.pm
%%SITE_PERL%%/Zonemaster/Backend/DB/PostgreSQL.pm
%%SITE_PERL%%/Zonemaster/Backend/DB/SQLite.pm
%%SITE_PERL%%/Zonemaster/Backend/Errors.pm
%%SITE_PERL%%/Zonemaster/Backend/Log.pm
%%SITE_PERL%%/Zonemaster/Backend/Metrics.pm
%%SITE_PERL%%/Zonemaster/Backend/RPCAPI.pm
%%SITE_PERL%%/Zonemaster/Backend/TestAgent.pm
%%SITE_PERL%%/Zonemaster/Backend/Translator.pm
%%SITE_PERL%%/Zonemaster/Backend/Validator.pm
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/GNUmakefile
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/backend_config.ini
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/cleanup-mysql.sql
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/cleanup-postgres.sql
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/create_db.pl
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/freebsd-pwd.conf
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/locale/da/LC_MESSAGES/Zonemaster-Backend.mo
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/locale/es/LC_MESSAGES/Zonemaster-Backend.mo
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/locale/fi/LC_MESSAGES/Zonemaster-Backend.mo
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/locale/fr/LC_MESSAGES/Zonemaster-Backend.mo
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/locale/nb/LC_MESSAGES/Zonemaster-Backend.mo
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/locale/sl/LC_MESSAGES/Zonemaster-Backend.mo
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/locale/sv/LC_MESSAGES/Zonemaster-Backend.mo
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/patch/README.txt
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/patch/patch_db_zonemaster_backend_ver_11.1.0.pl
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/patch/patch_db_zonemaster_backend_ver_11.2.0.pl
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/patch/patch_db_zonemaster_backend_ver_9.0.0.pl
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/patch/patch_mysql_db_zonemaster_backend_ver_1.0.3.pl
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/patch/patch_mysql_db_zonemaster_backend_ver_5.0.0.pl
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/patch/patch_mysql_db_zonemaster_backend_ver_5.0.2.pl
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/patch/patch_mysql_db_zonemaster_backend_ver_8.0.0.pl
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/patch/patch_postgresql_db_zonemaster_backend_ver_1.0.3.pl
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/patch/patch_postgresql_db_zonemaster_backend_ver_5.0.0.pl
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/patch/patch_postgresql_db_zonemaster_backend_ver_8.0.0.pl
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/patch/patch_sqlite_db_zonemaster_backend_ver_8.0.0.pl
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/tmpfiles.conf
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/zm-rpcapi.lsb
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/zm-rpcapi.service
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/zm-testagent.lsb
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/zm-testagent.service
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/zm_rpcapi-bsd
%%SITE_PERL%%/auto/share/dist/Zonemaster-Backend/zm_testagent-bsd
%%PERL5_MAN1%%/zmb.1.gz
%%PERL5_MAN1%%/zonemaster_backend_testagent.1.gz
%%PERL5_MAN3%%/Zonemaster::Backend.3.gz
%%PERL5_MAN3%%/Zonemaster::Backend::Config.3.gz
%%PERL5_MAN3%%/Zonemaster::Backend::Config::DCPlugin.3.gz
%%PERL5_MAN3%%/Zonemaster::Backend::DB.3.gz
%%PERL5_MAN3%%/Zonemaster::Backend::DB::MySQL.3.gz
%%PERL5_MAN3%%/Zonemaster::Backend::DB::PostgreSQL.3.gz
%%PERL5_MAN3%%/Zonemaster::Backend::DB::SQLite.3.gz
%%PERL5_MAN3%%/Zonemaster::Backend::Validator.3.gz
