# net/netatalk3/files/openrc-netatalk.in
#!/sbin/openrc-run

name=netatalk
command=%%PREFIX%%/sbin/netatalk

depend() {
	need %%ZEROCONF%%
	keyword -shutdown
}
