# net-mgmt/arpwatch/pkg-plist
sbin/arpwatch
sbin/arpsnmp
man/man8/arpwatch.8.gz
man/man8/arpsnmp.8.gz
arpwatch/arp2ethers
arpwatch/ethercodes.dat
arpwatch/d.awk
arpwatch/e.awk
arpwatch/p.awk
etc/rc.d/arpwatch.sh
