# net-mgmt/arpwatch/pkg-plist
arpwatch/arp2ethers
arpwatch/d.awk
arpwatch/duplicates.awk
arpwatch/e.awk
arpwatch/euppertolower.awk
arpwatch/massagevendor
arpwatch/p.awk
arpwatch/update-ethercodes
etc/rc.d/arpwatch
man/man8/arpsnmp.8.gz
man/man8/arpwatch.8.gz
sbin/arpsnmp
sbin/arpwatch
@touch arpwatch/arp.dat
@rmempty arpwatch/arp.dat
@dir(,operator,0755) arpwatch
