# net-mgmt/arpwatch/pkg-plist
arpwatch/arp2ethers
arpwatch/d.awk
arpwatch/duplicates.awk
arpwatch/e.awk
%%ETHERCODES%%arpwatch/ethercodes.dat
arpwatch/euppertolower.awk
arpwatch/massagevendor
arpwatch/p.awk
etc/rc.d/arpwatch
man/man8/arpsnmp.8.gz
man/man8/arpwatch.8.gz
sbin/arpsnmp
sbin/arpwatch
@postexec test -f %D/arpwatch/arp.dat || touch %D/arpwatch/arp.dat
@preunexec test -f %D/arpwatch/arp.dat && test -s %D/arpwatch/arp.dat || rm -f %D/arpwatch/arp.dat
@dir(,operator,0755) arpwatch
