# comms/aprsd/pkg-plist
bin/aprsd
bin/aprspass
etc/aprsd/INIT.TNC.samp
etc/aprsd/RESTORE.TNC.samp
etc/aprsd/aprsd.conf.samp
etc/aprsd/user.deny.samp
etc/aprsd/welcome.txt.samp
etc/rc.d/aprsd
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/chkaprsd
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/udp_example
%%PORTDOCS%%%%DOCSDIR%%/aprsddoc.html
%%PORTDOCS%%%%DOCSDIR%%/ports.html
%%PORTDOCS%%%%DOCSDIR%%/q.html
%%PORTDOCS%%%%DOCSDIR%%/qalgorithm.html
%%PORTDOCS%%@dirrm %%DOCSDIR%%
%%PORTEXAMPLES%%@dirrm %%EXAMPLESDIR%%
@dirrmtry etc/aprsd
