# sysutils/puppet-devel/pkg-plist
@comment $FreeBSD$
bin/puppetd
bin/puppetca
bin/puppet
bin/puppetdoc
bin/puppetrun
bin/puppetmasterd
%%RUBY_SITELIBDIR%%/puppet.rb
%%RUBY_SITELIBDIR%%/puppet/transportable.rb
%%RUBY_SITELIBDIR%%/puppet/provider.rb
%%RUBY_SITELIBDIR%%/puppet/statechange.rb
%%RUBY_SITELIBDIR%%/puppet/configuration.rb
%%RUBY_SITELIBDIR%%/puppet/autoload.rb
%%RUBY_SITELIBDIR%%/puppet/error.rb
%%RUBY_SITELIBDIR%%/puppet/parameter.rb
%%RUBY_SITELIBDIR%%/puppet/config.rb
%%RUBY_SITELIBDIR%%/puppet/daemon.rb
%%RUBY_SITELIBDIR%%/puppet/server.rb
%%RUBY_SITELIBDIR%%/puppet/transaction.rb
%%RUBY_SITELIBDIR%%/puppet/loadedfile.rb
%%RUBY_SITELIBDIR%%/puppet/inifile.rb
%%RUBY_SITELIBDIR%%/puppet/rails.rb
%%RUBY_SITELIBDIR%%/puppet/pgraph.rb
%%RUBY_SITELIBDIR%%/puppet/dsl.rb
%%RUBY_SITELIBDIR%%/puppet/filetype.rb
%%RUBY_SITELIBDIR%%/puppet/relationship.rb
%%RUBY_SITELIBDIR%%/puppet/log.rb
%%RUBY_SITELIBDIR%%/puppet/event.rb
%%RUBY_SITELIBDIR%%/puppet/networkclient.rb
%%RUBY_SITELIBDIR%%/puppet/sslcertificates.rb
%%RUBY_SITELIBDIR%%/puppet/type.rb
%%RUBY_SITELIBDIR%%/puppet/storage.rb
%%RUBY_SITELIBDIR%%/puppet/client.rb
%%RUBY_SITELIBDIR%%/puppet/util.rb
%%RUBY_SITELIBDIR%%/puppet/element.rb
%%RUBY_SITELIBDIR%%/puppet/feature.rb
%%RUBY_SITELIBDIR%%/puppet/metric.rb
%%RUBY_SITELIBDIR%%/puppet/suidmanager.rb
%%RUBY_SITELIBDIR%%/puppet/type/parsedtype.rb
%%RUBY_SITELIBDIR%%/puppet/type/cron.rb
%%RUBY_SITELIBDIR%%/puppet/type/sshkey.rb
%%RUBY_SITELIBDIR%%/puppet/type/resources.rb
%%RUBY_SITELIBDIR%%/puppet/type/component.rb
%%RUBY_SITELIBDIR%%/puppet/type/pfilebucket.rb
%%RUBY_SITELIBDIR%%/puppet/type/group.rb
%%RUBY_SITELIBDIR%%/puppet/type/notify.rb
%%RUBY_SITELIBDIR%%/puppet/type/zone.rb
%%RUBY_SITELIBDIR%%/puppet/type/tidy.rb
%%RUBY_SITELIBDIR%%/puppet/type/exec.rb
%%RUBY_SITELIBDIR%%/puppet/type/mount.rb
%%RUBY_SITELIBDIR%%/puppet/type/yumrepo.rb
%%RUBY_SITELIBDIR%%/puppet/type/service.rb
%%RUBY_SITELIBDIR%%/puppet/type/state.rb
%%RUBY_SITELIBDIR%%/puppet/type/user.rb
%%RUBY_SITELIBDIR%%/puppet/type/host.rb
%%RUBY_SITELIBDIR%%/puppet/type/schedule.rb
%%RUBY_SITELIBDIR%%/puppet/type/package.rb
%%RUBY_SITELIBDIR%%/puppet/type/port.rb
%%RUBY_SITELIBDIR%%/puppet/type/pfile.rb
%%RUBY_SITELIBDIR%%/puppet/type/pfile/target.rb
%%RUBY_SITELIBDIR%%/puppet/type/pfile/group.rb
%%RUBY_SITELIBDIR%%/puppet/type/pfile/content.rb
%%RUBY_SITELIBDIR%%/puppet/type/pfile/checksum.rb
%%RUBY_SITELIBDIR%%/puppet/type/pfile/mode.rb
%%RUBY_SITELIBDIR%%/puppet/type/pfile/source.rb
%%RUBY_SITELIBDIR%%/puppet/type/pfile/type.rb
%%RUBY_SITELIBDIR%%/puppet/type/pfile/ensure.rb
%%RUBY_SITELIBDIR%%/puppet/type/pfile/owner.rb
%%RUBY_SITELIBDIR%%/puppet/util/warnings.rb
%%RUBY_SITELIBDIR%%/puppet/util/execution.rb
%%RUBY_SITELIBDIR%%/puppet/util/graph.rb
%%RUBY_SITELIBDIR%%/puppet/util/docs.rb
%%RUBY_SITELIBDIR%%/puppet/util/pidlock.rb
%%RUBY_SITELIBDIR%%/puppet/util/posix.rb
%%RUBY_SITELIBDIR%%/puppet/util/methodhelper.rb
%%RUBY_SITELIBDIR%%/puppet/util/logging.rb
%%RUBY_SITELIBDIR%%/puppet/util/fileparsing.rb
%%RUBY_SITELIBDIR%%/puppet/util/variables.rb
%%RUBY_SITELIBDIR%%/puppet/util/metaid.rb
%%RUBY_SITELIBDIR%%/puppet/util/package.rb
%%RUBY_SITELIBDIR%%/puppet/util/classgen.rb
%%RUBY_SITELIBDIR%%/puppet/util/plist.rb
%%RUBY_SITELIBDIR%%/puppet/util/errors.rb
%%RUBY_SITELIBDIR%%/puppet/util/plist/parser.rb
%%RUBY_SITELIBDIR%%/puppet/util/plist/generator.rb
%%RUBY_SITELIBDIR%%/puppet/util/rails/collection_merger.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr.rb
%%RUBY_SITELIBDIR%%/puppet/external/lock.rb
%%RUBY_SITELIBDIR%%/puppet/external/event-loop.rb
%%RUBY_SITELIBDIR%%/puppet/external/base64.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/rdot.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/graph.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/common.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/edge.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/undirected_graph.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/labels.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/comparability.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/graph_api.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/strong_components.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/dot.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/digraph_distance.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/search.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/base.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/digraph.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/maximum_flow.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/biconnected.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/adjacency_graph.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/import.rb
%%RUBY_SITELIBDIR%%/puppet/external/gratr/chinese_postman.rb
%%RUBY_SITELIBDIR%%/puppet/external/event-loop/better-definers.rb
%%RUBY_SITELIBDIR%%/puppet/external/event-loop/event-loop.rb
%%RUBY_SITELIBDIR%%/puppet/external/event-loop/signal-system.rb
%%RUBY_SITELIBDIR%%/puppet/metatype/providers.rb
%%RUBY_SITELIBDIR%%/puppet/metatype/evaluation.rb
%%RUBY_SITELIBDIR%%/puppet/metatype/manager.rb
%%RUBY_SITELIBDIR%%/puppet/metatype/tags.rb
%%RUBY_SITELIBDIR%%/puppet/metatype/attributes.rb
%%RUBY_SITELIBDIR%%/puppet/metatype/closure.rb
%%RUBY_SITELIBDIR%%/puppet/metatype/schedules.rb
%%RUBY_SITELIBDIR%%/puppet/metatype/instances.rb
%%RUBY_SITELIBDIR%%/puppet/metatype/metaparams.rb
%%RUBY_SITELIBDIR%%/puppet/metatype/container.rb
%%RUBY_SITELIBDIR%%/puppet/metatype/relationships.rb
%%RUBY_SITELIBDIR%%/puppet/reports/tagmail.rb
%%RUBY_SITELIBDIR%%/puppet/reports/rrdgraph.rb
%%RUBY_SITELIBDIR%%/puppet/reports/log.rb
%%RUBY_SITELIBDIR%%/puppet/reports/store.rb
%%RUBY_SITELIBDIR%%/puppet/rails/source_file.rb
%%RUBY_SITELIBDIR%%/puppet/rails/resource.rb
%%RUBY_SITELIBDIR%%/puppet/rails/fact_name.rb
%%RUBY_SITELIBDIR%%/puppet/rails/fact_value.rb
%%RUBY_SITELIBDIR%%/puppet/rails/param_name.rb
%%RUBY_SITELIBDIR%%/puppet/rails/host.rb
%%RUBY_SITELIBDIR%%/puppet/rails/param_value.rb
%%RUBY_SITELIBDIR%%/puppet/rails/puppet_class.rb
%%RUBY_SITELIBDIR%%/puppet/rails/external/tagging/tag.rb
%%RUBY_SITELIBDIR%%/puppet/rails/external/tagging/tagging.rb
%%RUBY_SITELIBDIR%%/puppet/rails/external/tagging/acts_as_taggable.rb
%%RUBY_SITELIBDIR%%/puppet/rails/external/tagging/init.rb
%%RUBY_SITELIBDIR%%/puppet/rails/database/schema.rb
%%RUBY_SITELIBDIR%%/puppet/feature/ldap.rb
%%RUBY_SITELIBDIR%%/puppet/feature/rails.rb
%%RUBY_SITELIBDIR%%/puppet/feature/usage.rb
%%RUBY_SITELIBDIR%%/puppet/client/ca.rb
%%RUBY_SITELIBDIR%%/puppet/client/reporter.rb
%%RUBY_SITELIBDIR%%/puppet/client/resource.rb
%%RUBY_SITELIBDIR%%/puppet/client/file.rb
%%RUBY_SITELIBDIR%%/puppet/client/master.rb
%%RUBY_SITELIBDIR%%/puppet/client/log.rb
%%RUBY_SITELIBDIR%%/puppet/client/runner.rb
%%RUBY_SITELIBDIR%%/puppet/client/dipper.rb
%%RUBY_SITELIBDIR%%/puppet/client/status.rb
%%RUBY_SITELIBDIR%%/puppet/client/proxy.rb
%%RUBY_SITELIBDIR%%/puppet/sslcertificates/ca.rb
%%RUBY_SITELIBDIR%%/puppet/sslcertificates/inventory.rb
%%RUBY_SITELIBDIR%%/puppet/sslcertificates/certificate.rb
%%RUBY_SITELIBDIR%%/puppet/parser/parser.rb
%%RUBY_SITELIBDIR%%/puppet/parser/functions.rb
%%RUBY_SITELIBDIR%%/puppet/parser/collector.rb
%%RUBY_SITELIBDIR%%/puppet/parser/resource.rb
%%RUBY_SITELIBDIR%%/puppet/parser/scope.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast.rb
%%RUBY_SITELIBDIR%%/puppet/parser/lexer.rb
%%RUBY_SITELIBDIR%%/puppet/parser/interpreter.rb
%%RUBY_SITELIBDIR%%/puppet/parser/templatewrapper.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/resourceref.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/resourceparam.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/hostclass.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/leaf.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/astarray.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/component.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/resourceoverride.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/function.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/else.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/vardef.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/selector.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/caseopt.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/collection.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/casestatement.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/node.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/ifstatement.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/resourcedefaults.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/tag.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/collexpr.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/branch.rb
%%RUBY_SITELIBDIR%%/puppet/parser/ast/resourcedef.rb
%%RUBY_SITELIBDIR%%/puppet/parser/resource/param.rb
%%RUBY_SITELIBDIR%%/puppet/parser/resource/reference.rb
%%RUBY_SITELIBDIR%%/puppet/transaction/report.rb
%%RUBY_SITELIBDIR%%/puppet/server/fileserver.rb
%%RUBY_SITELIBDIR%%/puppet/server/authconfig.rb
%%RUBY_SITELIBDIR%%/puppet/server/ca.rb
%%RUBY_SITELIBDIR%%/puppet/server/rights.rb
%%RUBY_SITELIBDIR%%/puppet/server/resource.rb
%%RUBY_SITELIBDIR%%/puppet/server/report.rb
%%RUBY_SITELIBDIR%%/puppet/server/logger.rb
%%RUBY_SITELIBDIR%%/puppet/server/master.rb
%%RUBY_SITELIBDIR%%/puppet/server/filebucket.rb
%%RUBY_SITELIBDIR%%/puppet/server/runner.rb
%%RUBY_SITELIBDIR%%/puppet/server/servlet.rb
%%RUBY_SITELIBDIR%%/puppet/server/authstore.rb
%%RUBY_SITELIBDIR%%/puppet/provider/nameservice.rb
%%RUBY_SITELIBDIR%%/puppet/provider/mount.rb
%%RUBY_SITELIBDIR%%/puppet/provider/parsedfile.rb
%%RUBY_SITELIBDIR%%/puppet/provider/host/parsed.rb
%%RUBY_SITELIBDIR%%/puppet/provider/host/netinfo.rb
%%RUBY_SITELIBDIR%%/puppet/provider/port/parsed.rb
%%RUBY_SITELIBDIR%%/puppet/provider/user/pw.rb
%%RUBY_SITELIBDIR%%/puppet/provider/user/useradd.rb
%%RUBY_SITELIBDIR%%/puppet/provider/user/netinfo.rb
%%RUBY_SITELIBDIR%%/puppet/provider/zone/solaris.rb
%%RUBY_SITELIBDIR%%/puppet/provider/group/pw.rb
%%RUBY_SITELIBDIR%%/puppet/provider/group/netinfo.rb
%%RUBY_SITELIBDIR%%/puppet/provider/group/groupadd.rb
%%RUBY_SITELIBDIR%%/puppet/provider/mount/parsed.rb
%%RUBY_SITELIBDIR%%/puppet/provider/mount/netinfo.rb
%%RUBY_SITELIBDIR%%/puppet/provider/service/debian.rb
%%RUBY_SITELIBDIR%%/puppet/provider/service/redhat.rb
%%RUBY_SITELIBDIR%%/puppet/provider/service/smf.rb
%%RUBY_SITELIBDIR%%/puppet/provider/service/gentoo.rb
%%RUBY_SITELIBDIR%%/puppet/provider/service/base.rb
%%RUBY_SITELIBDIR%%/puppet/provider/service/init.rb
%%RUBY_SITELIBDIR%%/puppet/provider/nameservice/pw.rb
%%RUBY_SITELIBDIR%%/puppet/provider/nameservice/objectadd.rb
%%RUBY_SITELIBDIR%%/puppet/provider/nameservice/netinfo.rb
%%RUBY_SITELIBDIR%%/puppet/provider/sshkey/parsed.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/up2date.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/portage.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/dpkg.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/pkgdmg.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/freebsd.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/sunfreeware.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/ports.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/apt.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/apple.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/gem.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/aptitude.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/rpm.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/darwinport.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/sun.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/yum.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/openbsd.rb
%%RUBY_SITELIBDIR%%/puppet/provider/package/blastwave.rb
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/util/rails
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/util/plist
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/util
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/type/pfile
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/type
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/transaction
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/sslcertificates
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/server
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/reports
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/rails/external/tagging
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/rails/external
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/rails/database
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/rails
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/provider/zone
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/provider/user
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/provider/sshkey
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/provider/service
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/provider/port
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/provider/package
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/provider/nameservice
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/provider/mount
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/provider/host
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/provider/group
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/provider
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/parser/resource
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/parser/ast
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/parser
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/metatype
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/feature
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/external/gratr
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/external/event-loop
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/external
@dirrmtry %%RUBY_SITELIBDIR%%/puppet/client
@dirrmtry %%RUBY_SITELIBDIR%%/puppet
%%PORTDOCS%%%%DOCSDIR%%/CHANGELOG
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/LICENSE
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%@dirrm %%DOCSDIR%%
%%PORTDOCS%%%%EXAMPLESDIR%%/root/bin/sleeper
%%PORTDOCS%%%%EXAMPLESDIR%%/root/etc/init.d/sleeper
%%PORTDOCS%%%%EXAMPLESDIR%%/root/etc/puppet/puppetd.conf
%%PORTDOCS%%%%EXAMPLESDIR%%/root/etc/puppet/fileserver.conf
%%PORTDOCS%%%%EXAMPLESDIR%%/root/etc/puppet/puppetmasterd.conf
%%PORTDOCS%%%%EXAMPLESDIR%%/root/etc/puppet/namespaceauth.conf
%%PORTDOCS%%%%EXAMPLESDIR%%/root/etc/puppet/tagmail.conf
%%PORTDOCS%%%%EXAMPLESDIR%%/root/etc/otherfile
%%PORTDOCS%%%%EXAMPLESDIR%%/root/etc/configfile
%%PORTDOCS%%%%EXAMPLESDIR%%/root/etc/debian-passwd
%%PORTDOCS%%%%EXAMPLESDIR%%/root/etc/debian-syslog.conf
%%PORTDOCS%%%%EXAMPLESDIR%%/code/allatonce
%%PORTDOCS%%%%EXAMPLESDIR%%/code/head
%%PORTDOCS%%%%EXAMPLESDIR%%/code/selectors
%%PORTDOCS%%%%EXAMPLESDIR%%/code/execs
%%PORTDOCS%%%%EXAMPLESDIR%%/code/nodes
%%PORTDOCS%%%%EXAMPLESDIR%%/code/filerecursion
%%PORTDOCS%%%%EXAMPLESDIR%%/code/components
%%PORTDOCS%%%%EXAMPLESDIR%%/code/file.bl
%%PORTDOCS%%%%EXAMPLESDIR%%/code/groups
%%PORTDOCS%%%%EXAMPLESDIR%%/code/functions
%%PORTDOCS%%%%EXAMPLESDIR%%/code/fileparsing
%%PORTDOCS%%%%EXAMPLESDIR%%/code/svncommit
%%PORTDOCS%%%%EXAMPLESDIR%%/code/filedefaults
%%PORTDOCS%%%%EXAMPLESDIR%%/code/simpletests
%%PORTDOCS%%%%EXAMPLESDIR%%/code/assignments
%%PORTDOCS%%%%EXAMPLESDIR%%/code/relationships
%%PORTDOCS%%%%EXAMPLESDIR%%/code/importing
%%PORTDOCS%%%%EXAMPLESDIR%%/code/one
%%PORTDOCS%%@dirrmtry %%EXAMPLESDIR%%/root/bin
%%PORTDOCS%%@dirrmtry %%EXAMPLESDIR%%/root/etc/init.d
%%PORTDOCS%%@dirrmtry %%EXAMPLESDIR%%/root/etc/puppet
%%PORTDOCS%%@dirrmtry %%EXAMPLESDIR%%/root/etc
%%PORTDOCS%%@dirrmtry %%EXAMPLESDIR%%/root
%%PORTDOCS%%@dirrmtry %%EXAMPLESDIR%%/code
%%PORTDOCS%%@dirrmtry %%EXAMPLESDIR%%
etc/puppet/puppet.conf-dist
@dirrmtry etc/puppet/ssl
@dirrmtry etc/puppet/manifests
@dirrmtry etc/puppet
