# net/exabgp/files/exabgp.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: exabgp
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following line(s) to /etc/rc.conf to enable exabgp:
#
# exabgp_enable="YES"

. /etc/rc.subr

name=exabgp
rcvar=exabgp_enable

load_rc_config $name

exabgp_enable=${exabgp_enable:-"NO"}
exabgp_conf=${exabgp_conf:-"%%ETCDIR%%/exabgp.conf"}

pidfile=/var/run/${name}/${name}.pid

required_files=${exabgp_conf}

command="%%PYTHON_CMD%%"
command_args="-m exabgp --folder /usr/local/etc/exabgp ${exabgp_conf}"

configtest_cmd="exabgp_configtest"
start_precmd="exabgp_precmd"
reload_precmd="exabgp_configtest"
restart_precmd="exabgp_configtest"
reload_all_cmd="exabgp_reload_all"
reload_all_precmd="exabgp_configtest"
extra_commands="configtest reload reload_all"
sig_reload="USR1"

exabgp_reload_all()
{
    	echo "Reloading exabgp configuration and processes."
	kill -USR2 $rc_pid
}

exabgp_configtest()
{
	echo "Performing sanity check on ${name} configuration."
	result=$(env exabgp.log.configuration=true exabgp.daemon.daemonize=false exabgp.log.level=ERR exabgp.log.destination=stdout ${command} -m ${name} --folder /usr/local/etc/exabgp -t ${exabgp_conf} 2>&1)
	if [ ! -z "${result}" ]; then
		echo "${result}"
		err 1 "Failed to validate configuration file, aborting."
	fi
}

exabgp_precmd()
{
	exabgp_configtest
	install -d -o _bgpd -g _bgpd -m 755 /var/run/${name}
	install -d -o _bgpd -g _bgpd -m 755 /var/log/${name}
	rm -f ${pidfile}
}

run_rc_command "$1"
