# net/exabgp/files/exabgp.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: exabgp
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following line(s) to /etc/rc.conf to enable exabgp:
#
# exabgp_enable="YES"

. /etc/rc.subr

name=exabgp
rcvar=exabgp_enable

load_rc_config $name

exabgp_enable=${exabgp_enable:-"NO"}
exabgp_conf=${exabgp_conf:-"%%ETCDIR%%/exabgp.conf"}

pidfile=/var/run/${name}/${name}.pid

required_files=${exabgp_conf}

command="%%PYTHON_SITELIBDIR%%/${name}/application.py"
procname="%%PYTHON_CMD%%"

start_cmd="exabgp_start"
extra_commands="reload"

exabgp_start()
{
	if [ -z "$rc_fast" -a -n "$rc_pid" ]; then
		echo 1>&2 "${name} already running? (pid=$rc_pid)."
		return 1
	fi

	echo "Starting ${name}."

	install -d -o _bgpd -g _bgpd -m 755 /var/run/${name}
	install -d -o _bgpd -g _bgpd -m 755 /var/log/${name}
	rm -f ${pidfile}

	${procname} ${command} -c %%ETCDIR%% ${exabgp_conf}
}

run_rc_command "$1"
