# news/inn-stable/pkg-plist
bin/inews
bin/rnews
etc/rc.news
etc/rc.d/innd.sh
include/inndcomm.h
include/libinn.h
include/qio.h
lib/libinn.a
sbin/innd
sbin/inndstart
sbin/in.nnrpd
news/lib/overview.fmt
news/lib/parsecontrol
news/lib/send-ihave
news/lib/send-nntp
news/lib/send-uucp
news/lib/innlog.pl
news/lib/innshellvars
news/lib/innshellvars.pl
news/lib/innshellvars.tcl
news/lib/innshellvars.csh
news/lib/newsfeeds
news/lib/hosts.nntp
news/lib/hosts.nntp.nolimit
news/lib/nnrp.access
news/lib/passwd.nntp
news/lib/inn.conf
news/lib/moderators
news/lib/control.ctl
news/lib/expire.ctl
news/lib/nntpsend.ctl
news/lib/innwatch.ctl
news/lib/distrib.pats
news/lib/actsync.cfg
news/lib/actsync.ign
news/bin/control/checkgroups
news/bin/control/default
news/bin/control/ihave
news/bin/control/newgroup
news/bin/control/rmgroup
news/bin/control/sendme
news/bin/control/sendsys
news/bin/control/senduuname
news/bin/control/version
news/bin/control/docheckgroups
news/bin/control/startup_innd.pl
news/bin/control/filter_innd.pl
news/bin/control/filter_nnrpd.pl
news/bin/control/startup.tcl
news/bin/control/filter.tcl
news/bin/control/sample.control
news/bin/rnews/c7unbatch
news/bin/rnews/decode
news/bin/rnews/encode
news/bin/ctlinnd
news/bin/getlist
news/bin/innconfval
news/bin/archive
news/bin/batcher
news/bin/buffchan
news/bin/cvtbatch
news/bin/filechan
news/bin/nntpget
news/bin/innxmit
news/bin/innxbatch
news/bin/overchan
news/bin/shlock
news/bin/shrinkfile
news/bin/crosspost
news/bin/sendxbatches
news/bin/actsync
news/bin/actsyncd
news/bin/actmerge
news/bin/convdate
news/bin/expire
news/bin/expireover
news/bin/fastrm
news/bin/grephistory
news/bin/makeactive
news/bin/makehistory
news/bin/newsrequeue
news/bin/prunehistory
news/bin/writelog
news/bin/innmail
news/bin/news.daily
news/bin/scanlogs
news/bin/expirerm
news/bin/tally.control
news/bin/tally.unwanted
news/bin/nntpsend
news/bin/sendbatch
news/bin/inncheck
news/bin/innstat
news/bin/pgpverify
news/bin/innwatch
news/bin/scanspool
