# news/inn-stable/pkg-plist
etc/rc.d/innd.sh
news/bin/actmerge
news/bin/actsync
news/bin/actsyncd
news/bin/archive
news/bin/auth/passwd/auth_smb
news/bin/auth/passwd/ckpasswd
news/bin/auth/passwd/radius
news/bin/auth/resolv/domain
news/bin/auth/resolv/ident
news/bin/batcher
news/bin/buffchan
news/bin/cnfsheadconf
news/bin/cnfsstat
news/bin/control/checkgroups.pl
news/bin/control/ihave.pl
news/bin/control/newgroup.pl
news/bin/control/rmgroup.pl
news/bin/control/sendme.pl
news/bin/control/sendsys.pl
news/bin/control/senduuname.pl
news/bin/control/version.pl
news/bin/controlbatch
news/bin/controlchan
news/bin/convdate
news/bin/ctlinnd
news/bin/cvtbatch
news/bin/docheckgroups
news/bin/expire
news/bin/expireover
news/bin/expirerm
news/bin/fastrm
news/bin/filechan
news/bin/filter/INN.py
news/bin/filter/filter.tcl
news/bin/filter/filter_innd.pl
news/bin/filter/filter_innd.py
news/bin/filter/filter_nnrpd.pl
news/bin/filter/nnrpd_auth.pl
news/bin/filter/nnrpd_auth.py
news/bin/filter/startup.tcl
news/bin/filter/startup_innd.pl
news/bin/getlist
news/bin/gpgverify
news/bin/grephistory
news/bin/imapfeed
news/bin/inews
news/bin/inncheck
news/bin/innconfval
news/bin/innd
news/bin/inndf
news/bin/inndstart
news/bin/innfeed
news/bin/innmail
news/bin/innreport
news/bin/innstat
news/bin/innupgrade
news/bin/innwatch
news/bin/innxbatch
news/bin/innxmit
news/bin/mailpost
news/bin/makedbz
news/bin/makehistory
news/bin/mod-active
news/bin/news.daily
news/bin/news2mail
news/bin/ninpaths
news/bin/nnrpd
news/bin/nntpget
news/bin/nntpsend
news/bin/ovdb_init
news/bin/ovdb_monitor
news/bin/ovdb_server
news/bin/ovdb_stat
news/bin/overchan
news/bin/perl-nocem
news/bin/pgpverify
news/bin/procbatch
news/bin/prunehistory
news/bin/pullnews
news/bin/rc.news
news/bin/rnews
news/bin/rnews.libexec/c7unbatch
news/bin/rnews.libexec/decode
news/bin/rnews.libexec/encode
news/bin/rnews.libexec/gunbatch
news/bin/scanlogs
news/bin/scanspool
news/bin/send-ihave
news/bin/send-nntp
news/bin/send-uucp
news/bin/sendinpaths
news/bin/sendxbatches
news/bin/shlock
news/bin/shrinkfile
news/bin/signcontrol
news/bin/simpleftp
news/bin/sm
news/bin/startinnfeed
news/bin/tally.control
news/bin/tdx-util
news/bin/writelog
news/db/active
news/db/active.times
news/db/history
news/db/history.dir
%%SUB_WITHOUT_TAGGED_HASH%%news/db/history.hash
%%SUB_WITHOUT_TAGGED_HASH%%news/db/history.index
%%SUB_WITH_TAGGED_HASH%%news/db/history.pag
news/db/newsgroups
news/doc/CONTRIBUTORS
news/doc/GPL
news/doc/HACKING
news/doc/INSTALL
news/doc/LICENSE
news/doc/NEWS
news/doc/README
news/doc/README.snapshot
news/doc/TODO
news/doc/compliance-nntp
news/doc/config-design
news/doc/config-semantics
news/doc/config-syntax
news/doc/external-auth
news/doc/history
news/doc/hook-perl
news/doc/hook-python
news/doc/hook-tcl
news/doc/sample-control
@unexec if cmp -s %D/news/etc/actsync.cfg %D/news/etc/actsync.cfg.dist; then rm -f %D/news/etc/actsync.cfg ; fi
news/etc/actsync.cfg.dist
@unexec if cmp -s %D/news/etc/actsync.ign %D/news/etc/actsync.ign.dist; then rm -f %D/news/etc/actsync.ign ; fi
news/etc/actsync.ign.dist
@unexec if cmp -s %D/news/etc/buffindexed.conf %D/news/etc/buffindexed.conf.dist; then rm -f %D/news/etc/buffindexed.conf ; fi
news/etc/buffindexed.conf.dist
@unexec if cmp -s %D/news/etc/control.ctl %D/news/etc/control.ctl.dist; then rm -f %D/news/etc/control.ctl ; fi
news/etc/control.ctl.dist
@unexec if cmp -s %D/news/etc/cycbuff.conf %D/news/etc/cycbuff.conf.dist; then rm -f %D/news/etc/cycbuff.conf ; fi
news/etc/cycbuff.conf.dist
@unexec if cmp -s %D/news/etc/distrib.pats %D/news/etc/distrib.pats.dist; then rm -f %D/news/etc/distrib.pats ; fi
news/etc/distrib.pats.dist
@unexec if cmp -s %D/news/etc/expire.ctl %D/news/etc/expire.ctl.dist; then rm -f %D/news/etc/expire.ctl ; fi
news/etc/expire.ctl.dist
@unexec if cmp -s %D/news/etc/incoming.conf %D/news/etc/incoming.conf.dist; then rm -f %D/news/etc/incoming.conf ; fi
news/etc/incoming.conf.dist
@unexec if cmp -s %D/news/etc/inn.conf %D/news/etc/inn.conf.dist; then rm -f %D/news/etc/inn.conf ; fi
news/etc/inn.conf.dist
@unexec if cmp -s %D/news/etc/innfeed.conf %D/news/etc/innfeed.conf.dist; then rm -f %D/news/etc/innfeed.conf ; fi
news/etc/innfeed.conf.dist
@unexec if cmp -s %D/news/etc/innreport.conf %D/news/etc/innreport.conf.dist; then rm -f %D/news/etc/innreport.conf ; fi
news/etc/innreport.conf.dist
@unexec if cmp -s %D/news/etc/innwatch.ctl %D/news/etc/innwatch.ctl.dist; then rm -f %D/news/etc/innwatch.ctl ; fi
news/etc/innwatch.ctl.dist
@unexec if cmp -s %D/news/etc/moderators %D/news/etc/moderators.dist; then rm -f %D/news/etc/moderators ; fi
news/etc/moderators.dist
@unexec if cmp -s %D/news/etc/motd.news %D/news/etc/motd.news.dist; then rm -f %D/news/etc/motd.news ; fi
news/etc/motd.news.dist
@unexec if cmp -s %D/news/etc/news2mail.cf %D/news/etc/news2mail.cf.dist; then rm -f %D/news/etc/news2mail.cf ; fi
news/etc/news2mail.cf.dist
@unexec if cmp -s %D/news/etc/newsfeeds %D/news/etc/newsfeeds.dist; then rm -f %D/news/etc/newsfeeds ; fi
news/etc/newsfeeds.dist
@unexec if cmp -s %D/news/etc/nnrpd.track %D/news/etc/nnrpd.track.dist; then rm -f %D/news/etc/nnrpd.track ; fi
news/etc/nnrpd.track.dist
@unexec if cmp -s %D/news/etc/nntpsend.ctl %D/news/etc/nntpsend.ctl.dist; then rm -f %D/news/etc/nntpsend.ctl ; fi
news/etc/nntpsend.ctl.dist
@unexec if cmp -s %D/news/etc/ovdb.conf %D/news/etc/ovdb.conf.dist; then rm -f %D/news/etc/ovdb.conf ; fi
news/etc/ovdb.conf.dist
@unexec if cmp -s %D/news/etc/overview.fmt %D/news/etc/overview.fmt.dist; then rm -f %D/news/etc/overview.fmt ; fi
news/etc/overview.fmt.dist
@unexec if cmp -s %D/news/etc/passwd.nntp %D/news/etc/passwd.nntp.dist; then rm -f %D/news/etc/passwd.nntp ; fi
news/etc/passwd.nntp.dist
@unexec if cmp -s %D/news/etc/radius.conf %D/news/etc/radius.conf.dist; then rm -f %D/news/etc/radius.conf ; fi 
news/etc/radius.conf.dist
@unexec if cmp -s %D/news/etc/readers.conf %D/news/etc/readers.conf.dist; then rm -f %D/news/etc/readers.conf; fi
news/etc/readers.conf.dist
@unexec if cmp -s %D/news/etc/sasl.conf %D/news/etc/sasl.conf.dist; then rm -f %D/news/etc/sasl.conf ; fi
news/etc/sasl.conf.dist
@unexec if cmp -s %D/news/etc/storage.conf %D/news/etc/storage.conf.dist; then rm -f %D/news/etc/storage.conf ; fi
news/etc/storage.conf.dist
@unexec if cmp -s %D/news/etc/subscriptions %D/news/etc/subscriptions.dist; then rm -f %D/news/etc/subscriptions ; fi
news/etc/subscriptions.dist
news/include/clibrary.h
news/include/conffile.h
news/include/config.h
news/include/dbz.h
news/include/inn/buffer.h
news/include/inn/confparse.h
news/include/inn/defines.h
news/include/inn/hashtab.h
news/include/inn/history.h
news/include/inn/innconf.h
news/include/inn/list.h
news/include/inn/md5.h
news/include/inn/messages.h
news/include/inn/mmap.h
news/include/inn/qio.h
news/include/inn/sequence.h
news/include/inn/system.h
news/include/inn/timer.h
news/include/inn/tst.h
news/include/inn/vector.h
news/include/inn/version.h
news/include/inn/wire.h
news/include/inndcomm.h
news/include/libinn.h
news/include/nntp.h
news/include/ov.h
news/include/paths.h
news/include/storage.h
news/lib/innreport_inn.pm
news/lib/innshellvars
news/lib/innshellvars.pl
news/lib/innshellvars.tcl
news/lib/libinn.a
news/lib/libinnhist.a
news/lib/libstorage.a
%%PORTDOCS%%%%DOCSDIR%%/CONTRIBUTORS
%%PORTDOCS%%%%DOCSDIR%%/HACKING
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/LICENSE
%%PORTDOCS%%%%DOCSDIR%%/MANIFEST
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@unexec rmdir %D/news/tmp 2>/dev/null || true
@dirrm news/spool/overview
@dirrm news/spool/outgoing
@dirrm news/spool/innfeed
@dirrm news/spool/incoming/bad
@dirrm news/spool/incoming
@dirrm news/spool/articles
@dirrm news/spool/archive
@dirrm news/spool/tmp
@dirrm news/spool
@dirrm news/run
@dirrm news/lib
@dirrm news/include/inn
@dirrm news/include
@unexec rmdir %D/news/etc 2>/dev/null || true
@dirrm news/doc
@unexec rmdir %D/news/db 2>/dev/null || true
@dirrm news/bin/rnews.libexec
@dirrm news/bin/filter
@dirrm news/bin/control
@dirrm news/bin/auth/resolv
@dirrm news/bin/auth/passwd
@dirrm news/bin/auth
@dirrm news/bin
@unexec rmdir %D/news 2>/dev/null ||( echo "===>  If you don't plan to use inn later, you can safely remove %D/news. Deinstall process kept these files intact:" | fmt && for i in $(find %D/news -type f); do echo "    - $i"; done)
