# net/gateway6/files/gateway6.in
#!/bin/sh
#
# $FreeBSD$
#

#
# PROVIDE: gateway6
# REQUIRE: FILESYSTEMS NETWORKING
# KEYWORD: shutdown
#

#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable the IPv6 tunnel to Freenet6.net:
#
# gateway6_enable (bool): Set to NO by default.
#                         Set it to YES to enable gateway6.
# gateway6_flags (flags): Set to "" by default.
#
# See gw6c(8) for gateway6_flags
#

. /etc/rc.subr

name="gateway6"
rcvar=gateway6_enable

command=%%PREFIX%%/bin/gw6c
required_files=%%PREFIX%%/etc/gw6c.conf

load_rc_config $name

: ${gateway6_enable="NO"}

run_rc_command "$1"
