# www/onlyoffice-documentserver/files/ds-converter.conf
[program:converter]
command=%%LOCALBASE%%/www/onlyoffice/documentserver/server/FileConverter/converter
directory=%%LOCALBASE%%/www/onlyoffice/documentserver/server/FileConverter
user=onlyoffice
environment=NODE_ENV=production-freebsd,NODE_CONFIG_DIR=%%LOCALBASE%%/etc/onlyoffice/documentserver,NODE_DISABLE_COLORS=1
stdout_logfile=/var/log/onlyoffice/documentserver/converter/out.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/onlyoffice/documentserver/converter/err.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
