# www/onlyoffice-documentserver/files/ds-docservice.conf
[program:docservice]
command=%%LOCALBASE%%/www/onlyoffice/documentserver/server/DocService/docservice
directory=%%LOCALBASE%%/www/onlyoffice/documentserver/server/DocService
user=onlyoffice
environment=NODE_ENV=production-freebsd,NODE_CONFIG_DIR=%%LOCALBASE%%/etc/onlyoffice/documentserver,NODE_DISABLE_COLORS=1
stdout_logfile=/var/log/onlyoffice/documentserver/docservice/out.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/onlyoffice/documentserver/docservice/err.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
