# www/onlyoffice-documentserver/files/ds-metrics.conf
[program:metrics]
command=%%LOCALBASE%%/www/onlyoffice/documentserver/server/Metrics/metrics ./config/config.js
directory=%%LOCALBASE%%/www/onlyoffice/documentserver/server/Metrics
user=onlyoffice
environment=NODE_DISABLE_COLORS=1
stdout_logfile=/var/log/onlyoffice/documentserver/metrics/out.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/onlyoffice/documentserver/metrics/err.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
