# sysutils/healthd/pkg-plist
bin/healthd
etc/healthd.conf.sample
etc/rc.d/healthd.sh
