# comms/ltmdm/pkg-plist
@unexec sh %D/etc/rc.d/ltmdm.sh stop || true
etc/rc.d/ltmdm.sh
%%DATADIR%%/ltmdm.ko
@dirrm %%DATADIR%%
