# comms/ltmdm/pkg-plist
@unexec sh %D/etc/rc.d/ltmdm.sh stop || true
