# net/nylon/files/openrc-nylon.in
#!/sbin/openrc-run

name=nylon
command=%%PREFIX%%/bin/nylon
nylon_flags=${nylon_flags:-"-s"}

depend() {
	keyword -shutdown
}
