# misc/py-cinder/files/openrc-cinder-volume.in
#!/sbin/openrc-run

# cinder_volume_logdir (str):
#	Set it to chagge log directory
#	Default is "/var/log/cinder"
#
# cinder_volume_args (str):
#	Set it to change command line arguments.
#	Default is "--log-file ${cinder_volume_logdir}/cinder-volume.log"
#

name=cinder_volume

PATH=%%PREFIX%%/bin:%%PREFIX%%/sbin:$PATH

pidfile="/var/run/cinder-volume.pid"

: ${cinder_volume_logdir:="/var/log/cinder"}
: ${cinder_volume_args:="--log-file ${cinder_volume_logdir}/cinder-volume.log"}

command="%%PREFIX%%/bin/python2.7"
supervisor=supervise-daemon
command_args="cinder-volume ${cinder_volume_args}"

start_pre() {
	mkdir -p ${cinder_volume_logdir}
}
depend()
{
	keyword -shutdown
}
