# net-im/openfire/files/openrc-openfire.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# openfire_user (username): Set to openfire by default.
#               Set it to required username.
# openfire_group (group):   Set to openfire by default.
#               Set it to required group.
# openfire_libdir (path):   Set to %%DATADIR%%/lib by default.
#               Set it to java classes directory.
# openfire_home (path):     Set to %%DATADIR%% by default.
#               Set it to java home directory.
# openfire_javargs (args):  Set to -Xmx256M by default.
#               See java -h for available arguments.

name="openfire"

# Set defaults
: ${openfire_user:=${name}}
: ${openfire_group:=${name}}
: ${openfire_libdir:=%%DATADIR%%/lib}
: ${openfire_home:=%%DATADIR%%}
: ${openfire_javargs:='-Xmx256M'}

pidfile=/var/run/${name}.pid

required_files="%%ETCDIR%%/openfire.xml"
java_options=" 	-server -jar ${openfire_javargs} \
		-Dopenfire.lib.dir=${openfire_libdir} \
		-DopenfireHome=${openfire_home}"

java_command="	%%LOCALBASE%%/bin/java ${java_options} \
		%%DATADIR%%/lib/startup.jar"

command="${java_command}"

start_pre() {
        touch ${pidfile}
        chown ${openfire_user}:${openfire_group} ${pidfile}
}

depend() {
	use net
	keyword -shutdown
}
