# sysutils/beats7/files/heartbeat.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: heartbeat
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable heartbeat:
#
# heartbeat_enable (bool):	Set to YES to enable heartbeat
# 				Default: NO
# heartbeat_flags (str):	Extra flags passed to heartbeat
# heartbeat_conf (str):		heartbeat configuration file
#				Default: ${PREFIX}/etc/heartbeat.yml

. /etc/rc.subr

name="heartbeat"
rcvar=${name}_enable
load_rc_config $name

: ${heartbeat_enable:="NO"}
: ${heartbeat_config:="%%ETCDIR%%"}
: ${heartbeat_home:="%%DATADIR%%/heartbeat"}
: ${heartbeat_logs:="/var/log/beats"}
: ${heartbeat_data:="/var/db/beats/heartbeat"}


# daemon
start_precmd=heartbeat_prestart
command=/usr/sbin/daemon
pidfile="/var/run/${name}"
command_args="-rP ${pidfile} %%PREFIX%%/sbin/${name} ${heartbeat_flags} -path.config ${heartbeat_config} --path.home ${heartbeat_home} --path.data ${heartbeat_data} --path.logs ${heartbeat_logs}"

heartbeat_prestart() {
# Have to empty rc_flags so they don't get passed to daemon(8)
	rc_flags=""
}

run_rc_command "$1"
