# net/smbfs/pkg-plist
@cwd /
dev/nsmb0
%%KMODDIR%%/smbfs.ko
sbin/mount_smbfs
@cwd /usr/local
bin/smbutil
etc/nsmb.conf.sample
etc/rc.d/smbfs.sh.sample
lib/libsmb.a
man/man1/smbutil.1.gz
man/man8/mount_smbfs.8.gz
