# sysutils/daemontools/pkg-plist
bin/svscan
bin/svscanboot
bin/supervise
bin/svc
bin/svok
bin/svstat
bin/fghack
bin/pgrphack
bin/multilog
bin/tai64n
bin/tai64nlocal
bin/softlimit
bin/setuidgid
bin/envuidgid
bin/envdir
bin/setlock
bin/readproctitle
etc/rc.d/%%STARTUP_SCRIPT%%
