# net/avahi-app/pkg-plist
bin/avahi-browse
bin/avahi-browse-domains
bin/avahi-publish
bin/avahi-publish-address
bin/avahi-publish-service
bin/avahi-resolve
bin/avahi-resolve-address
bin/avahi-resolve-host-name
bin/avahi-set-host-name
@unexec if cmp -s %D/etc/avahi/avahi-daemon.conf %D/etc/avahi/avahi-daemon.conf.dist; then rm -f %D/etc/avahi/avahi-daemon.conf; fi
etc/avahi/avahi-daemon.conf.dist
@exec [ -f %B/avahi-daemon.conf ] || cp %B/%f %B/avahi-daemon.conf
etc/avahi/avahi-dnsconfd.action
@unexec if cmp -s %D/etc/avahi/hosts %D/etc/avahi/hosts.dist; then rm -f %D/etc/avahi/hosts; fi
etc/avahi/hosts
@exec [ -f %B/hosts ] || cp %B/%f %B/hosts
etc/avahi/services/sftp-ssh.service
etc/avahi/services/ssh.service
etc/dbus-1/system.d/avahi-dbus.conf
etc/rc.d/avahi-daemon
etc/rc.d/avahi-dnsconfd
include/avahi-client/client.h
include/avahi-client/lookup.h
include/avahi-client/publish.h
include/avahi-common/address.h
include/avahi-common/alternative.h
include/avahi-common/cdecl.h
include/avahi-common/defs.h
include/avahi-common/domain.h
include/avahi-common/error.h
include/avahi-common/gccmacro.h
include/avahi-common/llist.h
include/avahi-common/malloc.h
include/avahi-common/rlist.h
include/avahi-common/simple-watch.h
include/avahi-common/strlst.h
include/avahi-common/thread-watch.h
include/avahi-common/timeval.h
include/avahi-common/watch.h
include/avahi-compat-howl/corby/buffer.h
include/avahi-compat-howl/corby/channel.h
include/avahi-compat-howl/corby/corby.h
include/avahi-compat-howl/corby/message.h
include/avahi-compat-howl/corby/object.h
include/avahi-compat-howl/corby/orb.h
include/avahi-compat-howl/discovery/discovery.h
include/avahi-compat-howl/discovery/text_record.h
include/avahi-compat-howl/howl.h
include/avahi-compat-howl/rendezvous/rendezvous.h
include/avahi-compat-howl/rendezvous/text_record.h
include/avahi-compat-howl/salt/address.h
include/avahi-compat-howl/salt/debug.h
include/avahi-compat-howl/salt/interface.h
include/avahi-compat-howl/salt/platform.h
include/avahi-compat-howl/salt/salt.h
include/avahi-compat-howl/salt/signal.h
include/avahi-compat-howl/salt/socket.h
include/avahi-compat-howl/salt/time.h
include/avahi-core/core.h
include/avahi-core/log.h
include/avahi-core/lookup.h
include/avahi-core/publish.h
include/avahi-core/rr.h
include/avahi-glib/glib-malloc.h
include/avahi-glib/glib-watch.h
include/avahi-gobject/ga-client.h
include/avahi-gobject/ga-entry-group.h
include/avahi-gobject/ga-enums.h
include/avahi-gobject/ga-error.h
include/avahi-gobject/ga-record-browser.h
include/avahi-gobject/ga-service-browser.h
include/avahi-gobject/ga-service-resolver.h
lib/avahi/.keep
lib/libavahi-client.a
lib/libavahi-client.la
lib/libavahi-client.so
lib/libavahi-client.so.3
lib/libavahi-common.a
lib/libavahi-common.la
lib/libavahi-common.so
lib/libavahi-common.so.3
lib/libavahi-core.a
lib/libavahi-core.la
lib/libavahi-core.so
lib/libavahi-core.so.6
lib/libavahi-glib.a
lib/libavahi-glib.la
lib/libavahi-glib.so
lib/libavahi-glib.so.1
lib/libavahi-gobject.a
lib/libavahi-gobject.la
lib/libavahi-gobject.so
lib/libavahi-gobject.so.0
lib/libhowl.a
lib/libhowl.la
lib/libhowl.so
lib/libhowl.so.0
libdata/pkgconfig/avahi-client.pc
libdata/pkgconfig/avahi-compat-howl.pc
libdata/pkgconfig/avahi-core.pc
libdata/pkgconfig/avahi-glib.pc
libdata/pkgconfig/avahi-gobject.pc
sbin/avahi-daemon
sbin/avahi-dnsconfd
%%DATADIR%%/avahi-service.dtd
%%DATADIR%%/introspection/AddressResolver.introspect
%%DATADIR%%/introspection/DomainBrowser.introspect
%%DATADIR%%/introspection/EntryGroup.introspect
%%DATADIR%%/introspection/HostNameResolver.introspect
%%DATADIR%%/introspection/RecordBrowser.introspect
%%DATADIR%%/introspection/Server.introspect
%%DATADIR%%/introspection/ServiceBrowser.introspect
%%DATADIR%%/introspection/ServiceResolver.introspect
%%DATADIR%%/introspection/ServiceTypeBrowser.introspect
%%DATADIR%%/service-types
share/locale/ca/LC_MESSAGES/avahi.mo
share/locale/de/LC_MESSAGES/avahi.mo
share/locale/el/LC_MESSAGES/avahi.mo
share/locale/en_AU/LC_MESSAGES/avahi.mo
share/locale/en_CA/LC_MESSAGES/avahi.mo
share/locale/en_GB/LC_MESSAGES/avahi.mo
share/locale/en_NZ/LC_MESSAGES/avahi.mo
share/locale/es/LC_MESSAGES/avahi.mo
share/locale/fi/LC_MESSAGES/avahi.mo
share/locale/fr/LC_MESSAGES/avahi.mo
share/locale/hu/LC_MESSAGES/avahi.mo
share/locale/it/LC_MESSAGES/avahi.mo
share/locale/ms/LC_MESSAGES/avahi.mo
share/locale/pl/LC_MESSAGES/avahi.mo
share/locale/pt_BR/LC_MESSAGES/avahi.mo
share/locale/ro/LC_MESSAGES/avahi.mo
share/locale/ru/LC_MESSAGES/avahi.mo
share/locale/sr/LC_MESSAGES/avahi.mo
share/locale/sv/LC_MESSAGES/avahi.mo
share/locale/sr@latin/LC_MESSAGES/avahi.mo
share/locale/zh_CN/LC_MESSAGES/avahi.mo
@dirrmtry share/locale/sr@latin/LC_MESSAGES
@dirrmtry share/locale/sr@latin
@dirrmtry share/locale/en_NZ/LC_MESSAGES
@dirrmtry share/locale/en_NZ
@dirrm %%DATADIR%%/introspection
@dirrm %%DATADIR%%
@dirrmtry share/applications
@dirrm lib/avahi
@dirrm include/avahi-gobject
@dirrm include/avahi-glib
@dirrm include/avahi-core
@dirrm include/avahi-compat-howl/salt
@dirrm include/avahi-compat-howl/rendezvous
@dirrm include/avahi-compat-howl/discovery
@dirrm include/avahi-compat-howl/corby
@dirrm include/avahi-compat-howl
@dirrm include/avahi-common
@dirrm include/avahi-client
@dirrm etc/avahi/services
@dirrm etc/avahi
