# net-mgmt/torrus/pkg-plist
bin/torrus
etc/rc.d/torrus.sh
etc/torrus/conf/devdiscover-siteconfig.pl.sample
etc/torrus/conf/email-siteconfig.pl.sample
etc/torrus/conf/notify-siteconfig.pl.sample
etc/torrus/conf/snmptrap-siteconfig.pl.sample
etc/torrus/conf/torrus-siteconfig.pl.sample
etc/torrus/discovery/README
etc/torrus/xmlconfig/site-global.xml.sample
libexec/torrus/acledit
libexec/torrus/action_notify
libexec/torrus/action_printemail
libexec/torrus/action_snmptrap
libexec/torrus/action_snmpv1trap
libexec/torrus/bdbinfo
libexec/torrus/buildsearchdb
libexec/torrus/cleanup
libexec/torrus/clearcache
libexec/torrus/collector
libexec/torrus/compilexml
libexec/torrus/configinfo
libexec/torrus/configsnapshot
libexec/torrus/devdiscover
libexec/torrus/genddx
libexec/torrus/genlist
libexec/torrus/genreport
libexec/torrus/install_plugin
libexec/torrus/monitor
libexec/torrus/rrddir2xml
libexec/torrus/schedulerinfo
libexec/torrus/scripts/rrdup_notify.sh
libexec/torrus/scripts/xml/extract-skeleton.xsl
libexec/torrus/srvderive
libexec/torrus/ttproclist
%%DOCSDIR%%/devdoc/architecture.pod
%%DOCSDIR%%/devdoc/devdiscover.pod
%%DOCSDIR%%/devdoc/progstyle.pod
%%DOCSDIR%%/devdoc/reqs.0.0.pod
%%DOCSDIR%%/devdoc/reqs.0.1.pod
%%DOCSDIR%%/devdoc/torrus_roadmap.pod
%%DOCSDIR%%/devdoc/wd.distributed.pod
%%DOCSDIR%%/devdoc/wd.messaging.pod
%%DOCSDIR%%/devdoc/wd.monitor-escalation.pod
%%DOCSDIR%%/devdoc/wd.uptime-mon.pod
%%DOCSDIR%%/install.txt
%%DOCSDIR%%/reporting_setup.txt
%%DOCSDIR%%/rpnexpr.txt
%%DOCSDIR%%/rrfw_torrus_migration.txt
%%DOCSDIR%%/scalability.txt
%%DOCSDIR%%/snmpdiscovery.txt
%%DOCSDIR%%/stylingprofile.txt
%%DOCSDIR%%/userguide.txt
%%DOCSDIR%%/vendorsupport.txt
%%DOCSDIR%%/webintf.txt
%%DOCSDIR%%/xmlconfig.txt
%%EXAMPLESDIR%%/README
%%EXAMPLESDIR%%/onms.tmpl
%%EXAMPLESDIR%%/onmsInterfaces.sh
%%EXAMPLESDIR%%/setmonitor.xupdate.xml
%%EXAMPLESDIR%%/torrus-siteconfig.powerbook.pl
share/torrus/conf_defaults/devdiscover-config.pl
share/torrus/conf_defaults/initscript.conf
share/torrus/conf_defaults/instvars
share/torrus/conf_defaults/torrus-config.pl
share/torrus/conf_defaults/webmux.pl
share/torrus/conf_defaults/webmux2.pl
share/torrus/sup/dtd/snmp-discovery.dtd
share/torrus/sup/dtd/torrus-config.dtd
share/torrus/sup/mibs/RRDTOOL-SMI.txt
share/torrus/sup/mibs/TORRUS-MIB.txt
share/torrus/sup/styling/colornames.pl
share/torrus/sup/styling/rainbow-schema.pl
share/torrus/sup/styling/torrus-schema.pl
share/torrus/sup/webplain/explain-rrdgraph.html
share/torrus/sup/webplain/torrus-printer.css
share/torrus/sup/webplain/torrus-report.css
share/torrus/sup/webplain/torrus.css
share/torrus/templates/aclexport.xml
share/torrus/templates/adminfo.html
share/torrus/templates/default-chooser.html
share/torrus/templates/default-dir.html
share/torrus/templates/default-helptext.html
share/torrus/templates/default-login.html
share/torrus/templates/default-recursivedir.html
share/torrus/templates/default-rrd.html
share/torrus/templates/default-tset.html
share/torrus/templates/email-alarm.txt
share/torrus/templates/expanded-dir.html
share/torrus/templates/globalsearch.html
share/torrus/templates/html-incblocks.txt
share/torrus/templates/overview-subleaves.html
share/torrus/templates/report-index.html
share/torrus/templates/report-monthly.html
share/torrus/templates/report-serviceid.html
share/torrus/templates/report-yearly.html
share/torrus/templates/search.html
share/torrus/templates/tset-list.html
share/torrus/xmlconfig/cdef-collector-defs.xml
share/torrus/xmlconfig/defaults.xml
share/torrus/xmlconfig/examples/apc-ups.xml
share/torrus/xmlconfig/examples/ascend.max.xml
share/torrus/xmlconfig/examples/docsis-monitors.xml
share/torrus/xmlconfig/examples/generic-netsnmp.xml
share/torrus/xmlconfig/examples/hpux.xml
share/torrus/xmlconfig/examples/monitors.xml
share/torrus/xmlconfig/examples/multigraph.xml
share/torrus/xmlconfig/examples/rainbow-schema.xml
share/torrus/xmlconfig/examples/servers.data
share/torrus/xmlconfig/examples/servers.tmpl
share/torrus/xmlconfig/generic/collector-periods.xml
share/torrus/xmlconfig/generic/monitors.xml
share/torrus/xmlconfig/generic/rfc1628.ups.xml
share/torrus/xmlconfig/generic/rfc1697.rdbms.xml
share/torrus/xmlconfig/generic/rfc2662.adsl-line.xml
share/torrus/xmlconfig/generic/rfc2670.docsis-if.xml
share/torrus/xmlconfig/generic/rfc2790.host-resources.xml
share/torrus/xmlconfig/generic/rfc2863.if-mib.xml
share/torrus/xmlconfig/old/cisco-mac-accounting-example.xml
share/torrus/xmlconfig/old/cisco.generic.old-0.1.4.xml
share/torrus/xmlconfig/old/cisco.ios.mac-accounting-0.1.8.xml
share/torrus/xmlconfig/old/rfc1213.xml
share/torrus/xmlconfig/old/rfc2670.docsis-if.old.0.1.5d-20040224.xml
share/torrus/xmlconfig/old/rfc2670.docsis-if.old.1.0.4.xml
share/torrus/xmlconfig/old/rfc2863.if-mib.old-0.1.4.xml
share/torrus/xmlconfig/old/rfc2863.if-mib.old-0.1.7.xml
share/torrus/xmlconfig/old/snmp-defs.old-0.1.2.xml
share/torrus/xmlconfig/snmp-defs.xml
share/torrus/xmlconfig/vendor/alteon.xml
share/torrus/xmlconfig/vendor/apc.ups.xml
share/torrus/xmlconfig/vendor/apple.ae.xml
share/torrus/xmlconfig/vendor/arbor_e.xml
share/torrus/xmlconfig/vendor/ascend.max.xml
share/torrus/xmlconfig/vendor/atmel.xml
share/torrus/xmlconfig/vendor/betternetworks.xml
share/torrus/xmlconfig/vendor/cisco.firewall.xml
share/torrus/xmlconfig/vendor/cisco.generic.xml
share/torrus/xmlconfig/vendor/cisco.ios.docsis.xml
share/torrus/xmlconfig/vendor/cisco.ios.mac-accounting.xml
share/torrus/xmlconfig/vendor/cisco.ios.xml
share/torrus/xmlconfig/vendor/cisco.sce.xml
share/torrus/xmlconfig/vendor/cisco.vdsl-line.xml
share/torrus/xmlconfig/vendor/compaq.cim.xml
share/torrus/xmlconfig/vendor/empire.systemedge.ntregperf.xml
share/torrus/xmlconfig/vendor/empire.systemedge.xml
share/torrus/xmlconfig/vendor/f5.bigip.xml
share/torrus/xmlconfig/vendor/junos.xml
share/torrus/xmlconfig/vendor/hp.hpux.xml
share/torrus/xmlconfig/vendor/liebert.xml
share/torrus/xmlconfig/vendor/microsoft.windows.xml
share/torrus/xmlconfig/vendor/motorola.bsr.xml
share/torrus/xmlconfig/vendor/netapp.filer.xml
share/torrus/xmlconfig/vendor/netscreen.xml
share/torrus/xmlconfig/vendor/paradyne.xdsl.xml
share/torrus/xmlconfig/vendor/symmetricom.xml
share/torrus/xmlconfig/vendor/ucd.ucd-snmp.xml
%%SITE_PERL%%/Torrus/ACL.pm
%%SITE_PERL%%/Torrus/ACL/AuthLocalMD5.pm
%%SITE_PERL%%/Torrus/ACL/Edit.pm
%%SITE_PERL%%/Torrus/ACL/Export.pm
%%SITE_PERL%%/Torrus/ACL/Import.pm
%%SITE_PERL%%/Torrus/Apache2Handler.pm
%%SITE_PERL%%/Torrus/ApacheHandler.pm
%%SITE_PERL%%/Torrus/Collector.pm
%%SITE_PERL%%/Torrus/Collector/CDef.pm
%%SITE_PERL%%/Torrus/Collector/CDef_Params.pm
%%SITE_PERL%%/Torrus/Collector/ExtDBI.pm
%%SITE_PERL%%/Torrus/Collector/ExternalStorage.pm
%%SITE_PERL%%/Torrus/Collector/RRDStorage.pm
%%SITE_PERL%%/Torrus/Collector/SNMP.pm
%%SITE_PERL%%/Torrus/Collector/SNMP_Params.pm
%%SITE_PERL%%/Torrus/ConfigBuilder.pm
%%SITE_PERL%%/Torrus/ConfigTree.pm
%%SITE_PERL%%/Torrus/ConfigTree/Validator.pm
%%SITE_PERL%%/Torrus/ConfigTree/Writer.pm
%%SITE_PERL%%/Torrus/ConfigTree/XMLCompiler.pm
%%SITE_PERL%%/Torrus/DB.pm
%%SITE_PERL%%/Torrus/DataAccess.pm
%%SITE_PERL%%/Torrus/DevDiscover.pm
%%SITE_PERL%%/Torrus/DevDiscover/ATMEL.pm
%%SITE_PERL%%/Torrus/DevDiscover/AlliedTelesyn_PBC18.pm
%%SITE_PERL%%/Torrus/DevDiscover/Alteon.pm
%%SITE_PERL%%/Torrus/DevDiscover/Apple_AE.pm
%%SITE_PERL%%/Torrus/DevDiscover/Arbor_E.pm
%%SITE_PERL%%/Torrus/DevDiscover/AscendMax.pm
%%SITE_PERL%%/Torrus/DevDiscover/AxxessIT.pm
%%SITE_PERL%%/Torrus/DevDiscover/BetterNetworks.pm
%%SITE_PERL%%/Torrus/DevDiscover/CiscoCatOS.pm
%%SITE_PERL%%/Torrus/DevDiscover/CiscoFirewall.pm
%%SITE_PERL%%/Torrus/DevDiscover/CiscoGeneric.pm
%%SITE_PERL%%/Torrus/DevDiscover/CiscoIOS.pm
%%SITE_PERL%%/Torrus/DevDiscover/CiscoIOS_Docsis.pm
%%SITE_PERL%%/Torrus/DevDiscover/CiscoIOS_MacAccounting.pm
%%SITE_PERL%%/Torrus/DevDiscover/CiscoIOS_SAA.pm
%%SITE_PERL%%/Torrus/DevDiscover/CiscoSCE.pm
%%SITE_PERL%%/Torrus/DevDiscover/CiscoVDSL.pm
%%SITE_PERL%%/Torrus/DevDiscover/CompaqCIM.pm
%%SITE_PERL%%/Torrus/DevDiscover/EmpireSystemedge.pm
%%SITE_PERL%%/Torrus/DevDiscover/F5BigIp.pm
%%SITE_PERL%%/Torrus/DevDiscover/JunOS.pm
%%SITE_PERL%%/Torrus/DevDiscover/Liebert.pm
%%SITE_PERL%%/Torrus/DevDiscover/MicrosoftWindows.pm
%%SITE_PERL%%/Torrus/DevDiscover/MotorolaBSR.pm
%%SITE_PERL%%/Torrus/DevDiscover/NetApp.pm
%%SITE_PERL%%/Torrus/DevDiscover/NetScreen.pm
%%SITE_PERL%%/Torrus/DevDiscover/OracleDatabase.pm
%%SITE_PERL%%/Torrus/DevDiscover/Paradyne.pm
%%SITE_PERL%%/Torrus/DevDiscover/RFC1628_UPS_MIB.pm
%%SITE_PERL%%/Torrus/DevDiscover/RFC1657_BGP4_MIB.pm
%%SITE_PERL%%/Torrus/DevDiscover/RFC1697_RDBMS.pm
%%SITE_PERL%%/Torrus/DevDiscover/RFC2011_IP_MIB.pm
%%SITE_PERL%%/Torrus/DevDiscover/RFC2662_ADSL_LINE.pm
%%SITE_PERL%%/Torrus/DevDiscover/RFC2670_DOCS_IF.pm
%%SITE_PERL%%/Torrus/DevDiscover/RFC2737_ENTITY_MIB.pm
%%SITE_PERL%%/Torrus/DevDiscover/RFC2790_HOST_RESOURCES.pm
%%SITE_PERL%%/Torrus/DevDiscover/RFC2863_IF_MIB.pm
%%SITE_PERL%%/Torrus/DevDiscover/Symmetricom.pm
%%SITE_PERL%%/Torrus/DevDiscover/UcdSnmp.pm
%%SITE_PERL%%/Torrus/DevDiscover/Xylan.pm
%%SITE_PERL%%/Torrus/Log.pm
%%SITE_PERL%%/Torrus/Monitor.pm
%%SITE_PERL%%/Torrus/RPN.pm
%%SITE_PERL%%/Torrus/Renderer.pm
%%SITE_PERL%%/Torrus/Renderer/AdmInfo.pm
%%SITE_PERL%%/Torrus/Renderer/Frontpage.pm
%%SITE_PERL%%/Torrus/Renderer/HTML.pm
%%SITE_PERL%%/Torrus/Renderer/RRDtool.pm
%%SITE_PERL%%/Torrus/ReportGenerator/MonthlySrvUsage.pm
%%SITE_PERL%%/Torrus/ReportGenerator.pm
%%SITE_PERL%%/Torrus/ReportOutput/HTML.pm
%%SITE_PERL%%/Torrus/ReportOutput.pm
%%SITE_PERL%%/Torrus/Search.pm
%%SITE_PERL%%/Torrus/ServiceID.pm
%%SITE_PERL%%/Torrus/Scheduler.pm
%%SITE_PERL%%/Torrus/SchedulerInfo.pm
%%SITE_PERL%%/Torrus/SiteConfig.pm
%%SITE_PERL%%/Torrus/SQL/Reports.pm
%%SITE_PERL%%/Torrus/SQL/SrvExport.pm
%%SITE_PERL%%/Torrus/SQL.pm
%%SITE_PERL%%/Torrus/TimeStamp.pm
@exec mkdir %D/var/torrus 2>/dev/null || true
@exec mkdir %D/var/torrus/cache 2>/dev/null || true
@exec mkdir %D/var/torrus/db 2>/dev/null || true
@exec mkdir %D/var/torrus/reports 2>/dev/null || true
@exec mkdir %D/var/torrus/session_data 2>/dev/null || true
@exec mkdir %D/var/torrus/session_data/lock 2>/dev/null || true
@exec mkdir %D/var/torrus/session_data/store 2>/dev/null || true
@exec mkdir %D/etc/torrus/templates 2>/dev/null || true
@exec mkdir %D/share/torrus/plugins 2>/dev/null || true
@exec mkdir %D/share/torrus/plugins/devdiscover-config 2>/dev/null || true
@exec mkdir %D/share/torrus/plugins/torrus-config 2>/dev/null || true
@exec mkdir %D/share/torrus/plugins/wrapper 2>/dev/null || true
@exec mkdir %D/share/torrus/sup/webscripts 2>/dev/null || true
@dirrm var/torrus/session_data/store
@dirrm var/torrus/session_data/lock
@dirrm var/torrus/session_data
@dirrm var/torrus/reports
@dirrmtry var/torrus/db
@dirrm var/torrus/cache
@dirrmtry var/torrus
@dirrm share/torrus/xmlconfig/vendor
@dirrm share/torrus/xmlconfig/old
@dirrm share/torrus/xmlconfig/generic
@dirrm share/torrus/xmlconfig/examples
@dirrm share/torrus/xmlconfig
@dirrm share/torrus/templates
@dirrm share/torrus/sup/webscripts
@dirrm share/torrus/sup/webplain
@dirrm share/torrus/sup/styling
@dirrm share/torrus/sup/mibs
@dirrm share/torrus/sup/dtd
@dirrm share/torrus/sup
@dirrm share/torrus/plugins/wrapper
@dirrm share/torrus/plugins/torrus-config
@dirrm share/torrus/plugins/devdiscover-config
@dirrm share/torrus/plugins
@dirrm share/torrus/conf_defaults
@dirrm share/torrus
@dirrm %%EXAMPLESDIR%%
@dirrm %%DOCSDIR%%/devdoc
@dirrm %%DOCSDIR%%
@dirrm libexec/torrus/scripts/xml
@dirrm libexec/torrus/scripts
@dirrm libexec/torrus
@dirrmtry etc/torrus/xmlconfig
@dirrmtry etc/torrus/templates
@dirrmtry etc/torrus/discovery
@dirrmtry etc/torrus/conf
@dirrmtry etc/torrus
@dirrm %%SITE_PERL%%/Torrus/SQL
@dirrm %%SITE_PERL%%/Torrus/ReportOutput
@dirrm %%SITE_PERL%%/Torrus/ReportGenerator
@dirrm %%SITE_PERL%%/Torrus/Renderer
@dirrm %%SITE_PERL%%/Torrus/DevDiscover
@dirrm %%SITE_PERL%%/Torrus/ConfigTree
@dirrm %%SITE_PERL%%/Torrus/Collector
@dirrm %%SITE_PERL%%/Torrus/ACL
@dirrm %%SITE_PERL%%/Torrus
@exec mkdir -p /var/log/torrus
