# net-mgmt/torrus/pkg-plist
bin/torrus
etc/rc.d/torrus.sh
etc/torrus/conf/devdiscover-siteconfig.pl.sample
etc/torrus/conf/email-siteconfig.pl.sample
etc/torrus/conf/notify-siteconfig.pl.sample
etc/torrus/conf/snmptrap-siteconfig.pl.sample
etc/torrus/conf/torrus-siteconfig.pl.sample
etc/torrus/discovery/README
etc/torrus/xmlconfig/site-global.xml.sample
torrus/bin/acledit
torrus/bin/action_notify
torrus/bin/action_printemail
torrus/bin/action_snmptrap
torrus/bin/action_snmpv1trap
torrus/bin/cleanup
torrus/bin/clearcache
torrus/bin/collector
torrus/bin/compilexml
torrus/bin/configinfo
torrus/bin/configsnapshot
torrus/bin/devdiscover
torrus/bin/genddx
torrus/bin/genlist
torrus/bin/genreport
torrus/bin/install_plugin
torrus/bin/monitor
torrus/bin/rrddir2xml
torrus/bin/schedulerinfo
torrus/bin/ttproclist
torrus/conf_defaults/devdiscover-config.pl
torrus/conf_defaults/initscript.conf
torrus/conf_defaults/instvars
torrus/conf_defaults/torrus-config.pl
torrus/conf_defaults/webmux.pl
torrus/conf_defaults/webmux2.pl
torrus/doc/devdoc/architecture.pod
torrus/doc/devdoc/devdiscover.pod
torrus/doc/devdoc/progstyle.pod
torrus/doc/devdoc/reqs.0.0.pod
torrus/doc/devdoc/reqs.0.1.pod
torrus/doc/devdoc/torrus_roadmap.pod
torrus/doc/devdoc/wd.distributed.pod
torrus/doc/devdoc/wd.messaging.pod
torrus/doc/devdoc/wd.monitor-escalation.pod
torrus/doc/devdoc/wd.uptime-mon.pod
torrus/doc/install.txt
torrus/doc/reporting_setup.txt
torrus/doc/rpnexpr.txt
torrus/doc/rrfw_torrus_migration.txt
torrus/doc/scalability.txt
torrus/doc/snmpdiscovery.txt
torrus/doc/stylingprofile.txt
torrus/doc/userguide.txt
torrus/doc/vendorsupport.txt
torrus/doc/webintf.txt
torrus/doc/xmlconfig.txt
torrus/examples/README
torrus/examples/onms.tmpl
torrus/examples/onmsInterfaces.sh
torrus/examples/setmonitor.xupdate.xml
torrus/examples/torrus-siteconfig.powerbook.pl
torrus/perllib/Torrus/ACL.pm
torrus/perllib/Torrus/ACL/AuthLocalMD5.pm
torrus/perllib/Torrus/ACL/Edit.pm
torrus/perllib/Torrus/ACL/Export.pm
torrus/perllib/Torrus/ACL/Import.pm
torrus/perllib/Torrus/Apache2Handler.pm
torrus/perllib/Torrus/ApacheHandler.pm
torrus/perllib/Torrus/Collector.pm
torrus/perllib/Torrus/Collector/ExtDBI.pm
torrus/perllib/Torrus/Collector/ExternalStorage.pm
torrus/perllib/Torrus/Collector/RRDStorage.pm
torrus/perllib/Torrus/Collector/SNMP.pm
torrus/perllib/Torrus/Collector/SNMP_Params.pm
torrus/perllib/Torrus/ConfigBuilder.pm
torrus/perllib/Torrus/ConfigTree.pm
torrus/perllib/Torrus/ConfigTree/Validator.pm
torrus/perllib/Torrus/ConfigTree/Writer.pm
torrus/perllib/Torrus/ConfigTree/XMLCompiler.pm
torrus/perllib/Torrus/DB.pm
torrus/perllib/Torrus/DataAccess.pm
torrus/perllib/Torrus/DevDiscover.pm
torrus/perllib/Torrus/DevDiscover/ATMEL.pm
torrus/perllib/Torrus/DevDiscover/AlliedTelesyn_PBC18.pm
torrus/perllib/Torrus/DevDiscover/AscendMax.pm
torrus/perllib/Torrus/DevDiscover/AxxessIT.pm
torrus/perllib/Torrus/DevDiscover/BetterNetworks.pm
torrus/perllib/Torrus/DevDiscover/CiscoCatOS.pm
torrus/perllib/Torrus/DevDiscover/CiscoFirewall.pm
torrus/perllib/Torrus/DevDiscover/CiscoGeneric.pm
torrus/perllib/Torrus/DevDiscover/CiscoIOS.pm
torrus/perllib/Torrus/DevDiscover/CiscoIOS_Docsis.pm
torrus/perllib/Torrus/DevDiscover/CiscoIOS_MacAccounting.pm
torrus/perllib/Torrus/DevDiscover/CiscoIOS_SAA.pm
torrus/perllib/Torrus/DevDiscover/CiscoVDSL.pm
torrus/perllib/Torrus/DevDiscover/CompaqCIM.pm
torrus/perllib/Torrus/DevDiscover/EmpireSystemedge.pm
torrus/perllib/Torrus/DevDiscover/F5BigIp.pm
torrus/perllib/Torrus/DevDiscover/MicrosoftWindows.pm
torrus/perllib/Torrus/DevDiscover/NetApp.pm
torrus/perllib/Torrus/DevDiscover/NetScreen.pm
torrus/perllib/Torrus/DevDiscover/OracleDatabase.pm
torrus/perllib/Torrus/DevDiscover/Paradyne.pm
torrus/perllib/Torrus/DevDiscover/RFC1657_BGP4_MIB.pm
torrus/perllib/Torrus/DevDiscover/RFC1697_RDBMS.pm
torrus/perllib/Torrus/DevDiscover/RFC2011_IP_MIB.pm
torrus/perllib/Torrus/DevDiscover/RFC2662_ADSL_LINE.pm
torrus/perllib/Torrus/DevDiscover/RFC2670_DOCS_IF.pm
torrus/perllib/Torrus/DevDiscover/RFC2737_ENTITY_MIB.pm
torrus/perllib/Torrus/DevDiscover/RFC2790_HOST_RESOURCES.pm
torrus/perllib/Torrus/DevDiscover/RFC2863_IF_MIB.pm
torrus/perllib/Torrus/DevDiscover/UcdSnmp.pm
torrus/perllib/Torrus/DevDiscover/Xylan.pm
torrus/perllib/Torrus/Log.pm
torrus/perllib/Torrus/Monitor.pm
torrus/perllib/Torrus/RPN.pm
torrus/perllib/Torrus/Renderer.pm
torrus/perllib/Torrus/Renderer/AdmInfo.pm
torrus/perllib/Torrus/Renderer/Frontpage.pm
torrus/perllib/Torrus/Renderer/HTML.pm
torrus/perllib/Torrus/Renderer/RRDtool.pm
torrus/perllib/Torrus/ReportGenerator/MonthlySrvUsage.pm
torrus/perllib/Torrus/ReportGenerator.pm
torrus/perllib/Torrus/ReportOutput/HTML.pm
torrus/perllib/Torrus/ReportOutput.pm
torrus/perllib/Torrus/ServiceID.pm
torrus/perllib/Torrus/Scheduler.pm
torrus/perllib/Torrus/SchedulerInfo.pm
torrus/perllib/Torrus/SiteConfig.pm
torrus/perllib/Torrus/SQL/Reports.pm
torrus/perllib/Torrus/SQL/SrvExport.pm
torrus/perllib/Torrus/SQL.pm
torrus/perllib/Torrus/TimeStamp.pm
torrus/scripts/rrdup_notify.sh
torrus/scripts/xml/extract-skeleton.xsl
torrus/sup/dtd/snmp-discovery.dtd
torrus/sup/dtd/torrus-config.dtd
torrus/sup/mibs/RRDTOOL-SMI.txt
torrus/sup/mibs/TORRUS-MIB.txt
torrus/sup/styling/colornames.pl
torrus/sup/styling/rainbow-schema.pl
torrus/sup/styling/torrus-schema.pl
torrus/sup/webplain/explain-rrdgraph.html
torrus/sup/webplain/torrus-printer.css
torrus/sup/webplain/torrus-report.css
torrus/sup/webplain/torrus.css
torrus/templates/aclexport.xml
torrus/templates/adminfo.html
torrus/templates/default-chooser.html
torrus/templates/default-dir.html
torrus/templates/default-helptext.html
torrus/templates/default-login.html
torrus/templates/default-recursivedir.html
torrus/templates/default-rrd.html
torrus/templates/default-tset.html
torrus/templates/email-alarm.txt
torrus/templates/expanded-dir.html
torrus/templates/html-incblocks.txt
torrus/templates/overview-subleaves.html
torrus/templates/report-index.html
torrus/templates/report-monthly.html
torrus/templates/report-serviceid.html
torrus/templates/report-yearly.html
torrus/templates/tset-list.html
torrus/xmlconfig/defaults.xml
torrus/xmlconfig/examples/apc-ups.xml
torrus/xmlconfig/examples/ascend.max.xml
torrus/xmlconfig/examples/docsis-monitors.xml
torrus/xmlconfig/examples/flowscan.xml
torrus/xmlconfig/examples/generic-netsnmp.xml
torrus/xmlconfig/examples/hpux.xml
torrus/xmlconfig/examples/monitors.xml
torrus/xmlconfig/examples/multigraph.xml
torrus/xmlconfig/examples/rainbow-schema.xml
torrus/xmlconfig/examples/servers.data
torrus/xmlconfig/examples/servers.tmpl
torrus/xmlconfig/generic/collector-periods.xml
torrus/xmlconfig/generic/monitors.xml
torrus/xmlconfig/generic/rfc2662.adsl-line.xml
torrus/xmlconfig/generic/rfc2670.docsis-if.xml
torrus/xmlconfig/generic/rfc2790.host-resources.xml
torrus/xmlconfig/generic/rfc2863.if-mib.xml
torrus/xmlconfig/old/cisco-mac-accounting-example.xml
torrus/xmlconfig/old/cisco.generic.old-0.1.4.xml
torrus/xmlconfig/old/cisco.ios.mac-accounting-0.1.8.xml
torrus/xmlconfig/old/rfc1213.xml
torrus/xmlconfig/old/rfc2670.docsis-if.old.0.1.5d-20040224.xml
torrus/xmlconfig/old/rfc2863.if-mib.old-0.1.4.xml
torrus/xmlconfig/old/rfc2863.if-mib.old-0.1.7.xml
torrus/xmlconfig/old/snmp-defs.old-0.1.2.xml
torrus/xmlconfig/snmp-defs.xml
torrus/xmlconfig/vendor/apc.ups.xml
torrus/xmlconfig/vendor/ascend.max.xml
torrus/xmlconfig/vendor/atmel.xml
torrus/xmlconfig/vendor/betternetworks.xml
torrus/xmlconfig/vendor/cisco.firewall.xml
torrus/xmlconfig/vendor/cisco.generic.xml
torrus/xmlconfig/vendor/cisco.ios.docsis.xml
torrus/xmlconfig/vendor/cisco.ios.mac-accounting.xml
torrus/xmlconfig/vendor/cisco.ios.xml
torrus/xmlconfig/vendor/cisco.vdsl-line.xml
torrus/xmlconfig/vendor/compaq.cim.xml
torrus/xmlconfig/vendor/empire.systemedge.ntregperf.xml
torrus/xmlconfig/vendor/empire.systemedge.xml
torrus/xmlconfig/vendor/f5.bigip.xml
torrus/xmlconfig/vendor/hp.hpux.xml
torrus/xmlconfig/vendor/microsoft.windows.xml
torrus/xmlconfig/vendor/netapp.filer.xml
torrus/xmlconfig/vendor/netscreen.xml
torrus/xmlconfig/vendor/paradyne.xdsl.xml
torrus/xmlconfig/vendor/smokeping.xml
torrus/xmlconfig/vendor/ucd.ucd-snmp.xml
@exec mkdir %D/var/torrus 2>/dev/null || true
@exec mkdir %D/var/torrus/cache 2>/dev/null || true
@exec mkdir %D/var/torrus/db 2>/dev/null || true
@exec mkdir %D/var/torrus/reports 2>/dev/null || true
@exec mkdir %D/var/torrus/session_data 2>/dev/null || true
@exec mkdir %D/var/torrus/session_data/lock 2>/dev/null || true
@exec mkdir %D/var/torrus/session_data/store 2>/dev/null || true
@exec mkdir %D/etc/torrus/templates 2>/dev/null || true
@exec mkdir %D/torrus/plugins 2>/dev/null || true
@exec mkdir %D/torrus/plugins/devdiscover-config 2>/dev/null || true
@exec mkdir %D/torrus/plugins/torrus-config 2>/dev/null || true
@exec mkdir %D/torrus/plugins/wrapper 2>/dev/null || true
@exec mkdir %D/torrus/sup/webscripts 2>/dev/null || true
@dirrm var/torrus/session_data/store
@dirrm var/torrus/session_data/lock
@dirrm var/torrus/session_data
@dirrm var/torrus/reports
@dirrm var/torrus/db
@dirrm var/torrus/collector_rrd
@dirrm var/torrus/cache
@dirrm var/torrus
@dirrm var
@dirrm torrus/xmlconfig/vendor
@dirrm torrus/xmlconfig/old
@dirrm torrus/xmlconfig/generic
@dirrm torrus/xmlconfig/examples
@dirrm torrus/xmlconfig
@dirrm torrus/templates
@dirrm torrus/sup/webscripts
@dirrm torrus/sup/webplain
@dirrm torrus/sup/styling
@dirrm torrus/sup/mibs
@dirrm torrus/sup/dtd
@dirrm torrus/sup
@dirrm torrus/scripts/xml
@dirrm torrus/scripts
@dirrm torrus/plugins/wrapper
@dirrm torrus/plugins/torrus-config
@dirrm torrus/plugins/devdiscover-config
@dirrm torrus/plugins
@dirrm torrus/perllib/Torrus/SQL
@dirrm torrus/perllib/Torrus/ReportOutput
@dirrm torrus/perllib/Torrus/ReportGenerator
@dirrm torrus/perllib/Torrus/Renderer
@dirrm torrus/perllib/Torrus/DevDiscover
@dirrm torrus/perllib/Torrus/ConfigTree
@dirrm torrus/perllib/Torrus/Collector
@dirrm torrus/perllib/Torrus/ACL
@dirrm torrus/perllib/Torrus
@dirrm torrus/perllib
@dirrm torrus/examples
@dirrm torrus/doc/devdoc
@dirrm torrus/doc
@dirrm torrus/conf_defaults
@dirrm torrus/bin
@dirrm torrus
@dirrm etc/torrus/xmlconfig
@dirrm etc/torrus/templates
@dirrm etc/torrus/discovery
@dirrm etc/torrus/conf
@dirrm etc/torrus
