# net/sixxs-aiccu/files/sixxs-aiccu.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: sixxs-aiccu
# REQUIRE: NETWORKING ldconfig

#
# Add the following lines to /etc/rc.conf to enable sixxs-aiccu:
#
#sixxs_aiccu_enable="YES"
#

. /etc/rc.subr

name="sixxs_aiccu"
rcvar=`set_rcvar`

command="%%PREFIX%%/sbin/sixxs-aiccu"

load_rc_config ${name}

# set default
: ${sixxs_aiccu_enable="NO"}
: ${sixxs_aiccu_config="%%PREFIX%%/etc/aiccu.conf"}

command_args="${sixxs_aiccu_config}"
required_files="${sixxs_aiccu_config}"
start_cmd="${command} start $command_args"
stop_cmd="${command} stop $command_args"
brokers_cmd="${command} brokers"
tunnels_cmd="${command} tunnels $command_args"
test_cmd="${command} test $command_args"
autotest_cmd="${command} autotest $command_args"
license_cmd="${command} license"
version_cmd="${command} version"
extra_commands="brokers tunnels test autotest license version"

run_rc_command "$1"
