# net/dhcpd/files/dhcpd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: dhcpd
# REQUIRE: DAEMON
# BEFORE:  LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable dhcpd:
# dhcpd_enable="YES"
# dhcpd_flags="<set as needed>"

. /etc/rc.subr

name=dhcpd
rcvar=dhcpd_enable

load_rc_config $name

command="%%PREFIX%%/sbin/dhcpd"

dhcpd_enable=${dhcpd_enable:-"NO"}
dhcpd_flags=${dhcpd_flags:-""}

run_rc_command "$1"
