# lang/erlang/files/openrc-epmd.in
#!/sbin/openrc-run

name=epmd

command=%%PREFIX%%/bin/epmd
command_args="-daemon"

stop_pre() {
	$command -kill >/dev/null
}
