# www/adzapper/pkg-plist
etc/adzapper.conf.default
etc/rc.d/adzapper.sh.sample
sbin/adzapper
share/adzapper/README
share/adzapper/adzapper.py
share/adzapper/doc/CHANGELOG
share/adzapper/doc/COPYING
share/adzapper/doc/COPYING-Zaplet
share/adzapper/doc/PythonPowered.gif
share/adzapper/doc/PythonPoweredSmall.gif
share/adzapper/doc/THANKS
share/adzapper/doc/TODO
share/adzapper/doc/about-adzapper.html
share/adzapper/doc/adzapper-bolt4.gif
share/adzapper/doc/description.txt
share/adzapper/doc/index.html
share/adzapper/doc/install.html
share/adzapper/doc/sourceforge-description.txt
share/adzapper/doc/template.html
share/adzapper/doc/zaplet-0.9.dtd
share/adzapper/doc/zaplet-file-format.html
share/adzapper/doc/zaplets.html
share/adzapper/doc/zpowered.jpg
share/adzapper/dtml/PythonPowered.gif
share/adzapper/dtml/PythonPoweredSmall.gif
share/adzapper/dtml/about.html
share/adzapper/dtml/about1.html
share/adzapper/dtml/adzapper-bolt4.gif
share/adzapper/dtml/blank.html
share/adzapper/dtml/default.dtml
share/adzapper/dtml/dot_clear.gif
share/adzapper/dtml/header1.html
share/adzapper/dtml/index.html
share/adzapper/dtml/logo.html
share/adzapper/dtml/restart.dtml
share/adzapper/dtml/restarted.dtml
share/adzapper/dtml/section1.dtml
share/adzapper/dtml/test.dtml
share/adzapper/dtml/thanks
share/adzapper/dtml/zapletwizard.dtml
share/adzapper/src/adzapper/CVS/Entries
share/adzapper/src/adzapper/CVS/Repository
share/adzapper/src/adzapper/CVS/Root
share/adzapper/src/adzapper/__init__.py
share/adzapper/src/adzapper/__init__.pyc
share/adzapper/src/adzapper/adzapper_engine.py
share/adzapper/src/adzapper/adzapper_engine.pyc
share/adzapper/src/adzapper/blank_gif.py
share/adzapper/src/adzapper/blank_gif.pyc
share/adzapper/src/adzapper/platform_dependent.py
share/adzapper/src/adzapper/platform_dependent.pyc
share/adzapper/src/asyncore_fixes.py
share/adzapper/src/asyncore_fixes.pyc
share/adzapper/src/config_file/__init__.py
share/adzapper/src/config_file/__init__.pyc
share/adzapper/src/config_file/adzapper.conf.xml
share/adzapper/src/config_file/config_file-0.2.dtd
share/adzapper/src/config_file/foo.xml
share/adzapper/src/config_file/foo1.xml
share/adzapper/src/config_file/rw_config_file.py
share/adzapper/src/config_file/rw_config_file.pyc
share/adzapper/src/debug/CVS/Entries
share/adzapper/src/debug/CVS/Repository
share/adzapper/src/debug/CVS/Root
share/adzapper/src/debug/__init__.py
share/adzapper/src/debug/__init__.pyc
share/adzapper/src/debug/debug.py
share/adzapper/src/debug/debug.pyc
share/adzapper/src/http_proxy.py
share/adzapper/src/http_proxy.pyc
share/adzapper/src/http_server_fixes.py
share/adzapper/src/http_server_fixes.pyc
share/adzapper/src/medusa/CVS/Entries
share/adzapper/src/medusa/CVS/Repository
share/adzapper/src/medusa/CVS/Root
share/adzapper/src/medusa/__init__.py
share/adzapper/src/medusa/__init__.pyc
share/adzapper/src/medusa/asynchat.py
share/adzapper/src/medusa/asynchat.pyc
share/adzapper/src/medusa/asyncore.py
share/adzapper/src/medusa/asyncore.pyc
share/adzapper/src/medusa/counter.py
share/adzapper/src/medusa/counter.pyc
share/adzapper/src/medusa/default_handler.py
share/adzapper/src/medusa/default_handler.pyc
share/adzapper/src/medusa/filesys.py
share/adzapper/src/medusa/filesys.pyc
share/adzapper/src/medusa/ftp_server.py
share/adzapper/src/medusa/ftp_server.pyc
share/adzapper/src/medusa/http_date.py
share/adzapper/src/medusa/http_date.pyc
share/adzapper/src/medusa/http_server.py
share/adzapper/src/medusa/http_server.pyc
share/adzapper/src/medusa/logger.py
share/adzapper/src/medusa/logger.pyc
share/adzapper/src/medusa/m_syslog.py
share/adzapper/src/medusa/m_syslog.pyc
share/adzapper/src/medusa/max_sockets.py
share/adzapper/src/medusa/max_sockets.pyc
share/adzapper/src/medusa/medusa_gif.py
share/adzapper/src/medusa/medusa_gif.pyc
share/adzapper/src/medusa/mime_type_table.py
share/adzapper/src/medusa/mime_type_table.pyc
share/adzapper/src/medusa/monitor.py
share/adzapper/src/medusa/monitor.pyc
share/adzapper/src/medusa/monitor_client.py
share/adzapper/src/medusa/monitor_client.pyc
share/adzapper/src/medusa/monitor_client_win32.py
share/adzapper/src/medusa/monitor_client_win32.pyc
share/adzapper/src/medusa/producers.py
share/adzapper/src/medusa/producers.pyc
share/adzapper/src/medusa/resolver.py
share/adzapper/src/medusa/resolver.pyc
share/adzapper/src/medusa/select_trigger.py
share/adzapper/src/medusa/select_trigger.pyc
share/adzapper/src/medusa/status_handler.py
share/adzapper/src/medusa/status_handler.pyc
share/adzapper/src/resolver_fixes.py
share/adzapper/src/resolver_fixes.pyc
share/adzapper/src/ui/__init__.py
share/adzapper/src/ui/__init__.pyc
share/adzapper/src/ui/actions.py
share/adzapper/src/ui/actions.pyc
share/adzapper/src/ui/dtml.py
share/adzapper/src/ui/dtml.pyc
share/adzapper/src/ui/html_table.py
share/adzapper/src/ui/html_table.pyc
share/adzapper/src/ui_handler.py
share/adzapper/src/ui_handler.pyc
share/adzapper/src/xml/__init__.py
share/adzapper/src/xml/__init__.pyc
share/adzapper/src/xml/parsers/__init__.py
share/adzapper/src/xml/parsers/__init__.pyc
share/adzapper/src/xml/parsers/sgmllib.py
share/adzapper/src/xml/parsers/sgmllib.pyc
share/adzapper/src/xml/parsers/xmllib.py
share/adzapper/src/xml/parsers/xmllib.pyc
share/adzapper/src/xml/parsers/xmlproc/__init__.py
share/adzapper/src/xml/parsers/xmlproc/__init__.pyc
share/adzapper/src/xml/parsers/xmlproc/catalog.py
share/adzapper/src/xml/parsers/xmlproc/catalog.pyc
share/adzapper/src/xml/parsers/xmlproc/charconv.py
share/adzapper/src/xml/parsers/xmlproc/charconv.pyc
share/adzapper/src/xml/parsers/xmlproc/dtdparser.py
share/adzapper/src/xml/parsers/xmlproc/dtdparser.pyc
share/adzapper/src/xml/parsers/xmlproc/errors.py
share/adzapper/src/xml/parsers/xmlproc/errors.pyc
share/adzapper/src/xml/parsers/xmlproc/namespace.py
share/adzapper/src/xml/parsers/xmlproc/namespace.pyc
share/adzapper/src/xml/parsers/xmlproc/xcatalog.py
share/adzapper/src/xml/parsers/xmlproc/xcatalog.pyc
share/adzapper/src/xml/parsers/xmlproc/xmlapp.py
share/adzapper/src/xml/parsers/xmlproc/xmlapp.pyc
share/adzapper/src/xml/parsers/xmlproc/xmldtd.py
share/adzapper/src/xml/parsers/xmlproc/xmldtd.pyc
share/adzapper/src/xml/parsers/xmlproc/xmlproc.py
share/adzapper/src/xml/parsers/xmlproc/xmlproc.pyc
share/adzapper/src/xml/parsers/xmlproc/xmlutils.py
share/adzapper/src/xml/parsers/xmlproc/xmlutils.pyc
share/adzapper/src/xml/parsers/xmlproc/xmlval.py
share/adzapper/src/xml/parsers/xmlproc/xmlval.pyc
share/adzapper/src/xml/sax/__init__.py
share/adzapper/src/xml/sax/__init__.pyc
share/adzapper/src/xml/sax/drivers/__init__.py
share/adzapper/src/xml/sax/drivers/__init__.pyc
share/adzapper/src/xml/sax/drivers/drv_xmllib.py
share/adzapper/src/xml/sax/drivers/drv_xmllib.pyc
share/adzapper/src/xml/sax/drivers/drv_xmlproc.py
share/adzapper/src/xml/sax/drivers/drv_xmlproc.pyc
share/adzapper/src/xml/sax/drivers/drv_xmlproc_val.py
share/adzapper/src/xml/sax/drivers/drv_xmlproc_val.pyc
share/adzapper/src/xml/sax/drivers/pylibs.py
share/adzapper/src/xml/sax/drivers/pylibs.pyc
share/adzapper/src/xml/sax/saxexts.py
share/adzapper/src/xml/sax/saxexts.pyc
share/adzapper/src/xml/sax/saxlib.py
share/adzapper/src/xml/sax/saxlib.pyc
share/adzapper/src/xml/sax/saxutils.py
share/adzapper/src/xml/sax/saxutils.pyc
share/adzapper/src/zaplet/CVS/Entries
share/adzapper/src/zaplet/CVS/Repository
share/adzapper/src/zaplet/CVS/Root
share/adzapper/src/zaplet/__init__.py
share/adzapper/src/zaplet/__init__.pyc
share/adzapper/src/zaplet/zaplet_engine.py
share/adzapper/src/zaplet/zaplet_engine.pyc
share/adzapper/src/zaplet/zaplet_object.py
share/adzapper/src/zaplet/zaplet_object.pyc
share/adzapper/src/zaplet/zaplet_sax_handler.py
share/adzapper/src/zaplet/zaplet_sax_handler.pyc
share/adzapper/zaplets/199.172.144.25.zap
share/adzapper/zaplets/206.41.20.7.zap
share/adzapper/zaplets/207.87.18.zap
share/adzapper/zaplets/209.207.224.zap
share/adzapper/zaplets/209.244.186.44.zap
share/adzapper/zaplets/216.207.47.35.zap
share/adzapper/zaplets/216.34.69.zap
share/adzapper/zaplets/216.35.211.245.zap
share/adzapper/zaplets/216.87.208.127.zap
share/adzapper/zaplets/247media.com.zap
share/adzapper/zaplets/ad.preferences.com.zap
share/adzapper/zaplets/ad.rouze.com.zap
share/adzapper/zaplets/adfu.blockstackers.com.zap
share/adzapper/zaplets/adimages.go.com.zap
share/adzapper/zaplets/admedia.xoom.com.zap
share/adzapper/zaplets/ads.freshmeat.net.zap
share/adzapper/zaplets/ads.imgis.com.zap
share/adzapper/zaplets/ads.netgravity.net.zap
share/adzapper/zaplets/ads.tripod.com.zap
share/adzapper/zaplets/ads.tucows.com.zap
share/adzapper/zaplets/ads.web.aol.com.zap
share/adzapper/zaplets/ads.x10.com.zap
share/adzapper/zaplets/adserv.quality-channel.de.zap
share/adzapper/zaplets/adserver.zap
share/adzapper/zaplets/adtech.de.zap
share/adzapper/zaplets/akamai.net.zap
share/adzapper/zaplets/akamaitech.net.zap
share/adzapper/zaplets/amedia.techies.com.zap
share/adzapper/zaplets/andovernews.com.zap
share/adzapper/zaplets/angelfire.com.zap
share/adzapper/zaplets/arstechnica.com.zap
share/adzapper/zaplets/banner.orb.net.zap
share/adzapper/zaplets/banners.uswestdex.com.zap
share/adzapper/zaplets/beyond2000.com.zap
share/adzapper/zaplets/bfast.com.zap
share/adzapper/zaplets/biz.thomsoninvest.net.zap
share/adzapper/zaplets/blueninja.com.zap
share/adzapper/zaplets/burstnet.com.zap
share/adzapper/zaplets/businessweek.com.zap
share/adzapper/zaplets/cdnow.com.zap
share/adzapper/zaplets/clickhere.egroups.com.zap
share/adzapper/zaplets/cmpnet.com.zap
share/adzapper/zaplets/cnet.com.zap
share/adzapper/zaplets/cnn.com.zap
share/adzapper/zaplets/co.uk.zap
share/adzapper/zaplets/codecatalog.com.zap
share/adzapper/zaplets/coffeecup.com.zap
share/adzapper/zaplets/core.freshmeat.net.zap
share/adzapper/zaplets/cygnus.com.zap
share/adzapper/zaplets/dailynews.yahoo.com.zap
share/adzapper/zaplets/dailynews1.yahoo.com.zap
share/adzapper/zaplets/dailyradar.com.zap
share/adzapper/zaplets/davecentral.com.zap
share/adzapper/zaplets/default-numeric.zap
share/adzapper/zaplets/default.zap
share/adzapper/zaplets/deja.com.zap
share/adzapper/zaplets/dhs.org.zap
share/adzapper/zaplets/dir.yahoo.com.zap
share/adzapper/zaplets/dogpile.com.zap
share/adzapper/zaplets/doubleclick.net.zap
share/adzapper/zaplets/download.macromedia.com.zap
share/adzapper/zaplets/eads.com.zap
share/adzapper/zaplets/economist.com.zap
share/adzapper/zaplets/feedmag.com.zap
share/adzapper/zaplets/flycast.com.zap
share/adzapper/zaplets/focalink.com.zap
share/adzapper/zaplets/fool.com.zap
share/adzapper/zaplets/free-banners.com.zap
share/adzapper/zaplets/freshmeat.net.zap
share/adzapper/zaplets/gnu.org.zap
share/adzapper/zaplets/go2net.com.zap
share/adzapper/zaplets/google.com.zap.bak
share/adzapper/zaplets/headhunter.net.zap
share/adzapper/zaplets/heise.de.zap
share/adzapper/zaplets/herring.com.zap
share/adzapper/zaplets/hondavfr.com.zap
share/adzapper/zaplets/hotwired.com.zap
share/adzapper/zaplets/i33.com.zap
share/adzapper/zaplets/ibm.com.zap
share/adzapper/zaplets/idirect.tucows.com.zap
share/adzapper/zaplets/image.pathfinder.com.zap
share/adzapper/zaplets/images.whatshotnow.com.zap
share/adzapper/zaplets/imgis.com.zap
share/adzapper/zaplets/infoseek.go.com.zap
share/adzapper/zaplets/insight.com.zap
share/adzapper/zaplets/interadnet.com.zap
share/adzapper/zaplets/intervu.net.zap
share/adzapper/zaplets/link4ads.com.zap
share/adzapper/zaplets/linkexchange.com.zap
share/adzapper/zaplets/linuxdevices.com.zap
share/adzapper/zaplets/linuxjournal.com.zap
share/adzapper/zaplets/linuxpower.com.zap
share/adzapper/zaplets/linuxstart.com.zap
share/adzapper/zaplets/linuxtoday.com.zap
share/adzapper/zaplets/lwn.net.zap
share/adzapper/zaplets/lycos.com.zap
share/adzapper/zaplets/mercurycenter.com.zap
share/adzapper/zaplets/money.pagecount.com.zap
share/adzapper/zaplets/msnbc.com.zap
share/adzapper/zaplets/netg.vitaminshoppe.com.zap
share/adzapper/zaplets/newshub.com.zap
share/adzapper/zaplets/ngadcenter.net.zap
share/adzapper/zaplets/ngserve.pcworld.com.zap
share/adzapper/zaplets/pic.geocities.com.zap
share/adzapper/zaplets/promotions.yahoo.com.zap
share/adzapper/zaplets/q.com.zap
share/adzapper/zaplets/redherring.com.zap
share/adzapper/zaplets/rouze.com.zap
share/adzapper/zaplets/salon.com.zap
share/adzapper/zaplets/salonmagazine.com.zap
share/adzapper/zaplets/securityportal.com.zap
share/adzapper/zaplets/sightings.com.zap
share/adzapper/zaplets/spiegel.com.zap
share/adzapper/zaplets/static.admaximize.com.zap
share/adzapper/zaplets/store.yahoo.com.zap
share/adzapper/zaplets/technocrat.net.zap
share/adzapper/zaplets/techrepublic.com.zap
share/adzapper/zaplets/thechipmerchant.com.zap
share/adzapper/zaplets/theonion.com.zap
share/adzapper/zaplets/theregister.co.uk.zap
share/adzapper/zaplets/thestandard.com.zap
share/adzapper/zaplets/thinksecret.com.zap
share/adzapper/zaplets/uexpress.com.zap
share/adzapper/zaplets/valueclick.com.zap
share/adzapper/zaplets/view.avenuea.com.zap
share/adzapper/zaplets/virtumundo.com.zap
share/adzapper/zaplets/was.ads.link4ads.com.zap
share/adzapper/zaplets/websponsors.com.zap
share/adzapper/zaplets/wired.com.zap
share/adzapper/zaplets/worldforge.org.zap
share/adzapper/zaplets/wsj.com.zap
share/adzapper/zaplets/www.ad.tomshardware.com.zap
share/adzapper/zaplets/www.go.com.zap
share/adzapper/zaplets/www.rzeczpospolita.pl.zap
share/adzapper/zaplets/x10.com.zap
share/adzapper/zaplets/xb.xoom.com.zap
share/adzapper/zaplets/yahoo.co.uk.zap
share/adzapper/zaplets/yahoo.com.zap
share/adzapper/zaplets/yimg.com.zap
share/adzapper/zaplets/yoda.cybereps.com.zap
share/adzapper/zaplets/zdnet.com.zap
@dirrm share/adzapper/zaplets
@dirrm share/adzapper/src/zaplet/CVS
@dirrm share/adzapper/src/zaplet
@dirrm share/adzapper/src/xml/sax/drivers
@dirrm share/adzapper/src/xml/sax
@dirrm share/adzapper/src/xml/parsers/xmlproc
@dirrm share/adzapper/src/xml/parsers
@dirrm share/adzapper/src/xml
@dirrm share/adzapper/src/ui
@dirrm share/adzapper/src/medusa/CVS
@dirrm share/adzapper/src/medusa
@dirrm share/adzapper/src/debug/CVS
@dirrm share/adzapper/src/debug
@dirrm share/adzapper/src/config_file
@dirrm share/adzapper/src/adzapper/CVS
@dirrm share/adzapper/src/adzapper
@dirrm share/adzapper/src
@dirrm share/adzapper/dtml
@dirrm share/adzapper/doc
@dirrm share/adzapper
