# net/ntopng/files/ntopng.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: ntopng
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable ntop:
#
# ntopng_enable (bool):		Set it to "YES" to enable ntop
#				Default is "NO".
# ntopng_flags (flags):		Set extra flags to ntop
#				Default is "-d --use-syslog=daemon". see ntop(8).
#

. /etc/rc.subr

name=ntopng
rcvar=ntopng_enable
pid_path=/var/run/${name}.pid
httpdocs_path=%%DATADIR%%/httpdocs
scripts_path=%%DATADIR%%/scripts
callbacks_path=%%DATADIR%%/scripts/callbacks
user_path=%%PREFIX%%/etc/ntopng-users.conf

load_rc_config $name

: ${ntopng_enable="NO"}
: ${ntopng_flags="-e -G ${pid_path}"}

command=%%PREFIX%%/bin/ntopng
command_args="-e -G $pid_path -1 $httpdocs_path -2 $scripts_path -3 $callbacks_path

run_rc_command "$1"
