# net/tcserver/pkg-plist
sbin/elsa
etc/elsa.cfg
etc/rc.d/elsa.sh
share/elsa/LICENSE
share/elsa/README
share/elsa/content/elsa-test.mov
@dirrm share/elsa/content
@dirrm share/elsa
