# net/tcserver/pkg-plist
sbin/tcserver
etc/tcserver.cfg
etc/tcserver.cfg-full
etc/rc.d/tcserver.sh
share/tcserver/LICENSE
share/tcserver/README
share/tcserver/TeraCASTAdmin11.pdf
share/tcserver/content/tcserver-test.mov
@dirrm share/tcserver/content
@dirrm share/tcserver
