# x11/xapp/pkg-plist
bin/pastebin
bin/upload-system-info
bin/xfce4-set-wallpaper
etc/X11/xinit/xinitrc.d/80xapp-gtk3-module.sh
etc/xdg/autostart/xapp-sn-watcher.desktop
include/xapp/libxapp/xapp-favorites.h
include/xapp/libxapp/xapp-gtk-window.h
include/xapp/libxapp/xapp-icon-chooser-button.h
include/xapp/libxapp/xapp-icon-chooser-dialog.h
include/xapp/libxapp/xapp-kbd-layout-controller.h
include/xapp/libxapp/xapp-monitor-blanker.h
include/xapp/libxapp/xapp-preferences-window.h
include/xapp/libxapp/xapp-stack-sidebar.h
include/xapp/libxapp/xapp-status-icon-monitor.h
include/xapp/libxapp/xapp-status-icon.h
include/xapp/libxapp/xapp-statusicon-interface.h
include/xapp/libxapp/xapp-style-manager.h
include/xapp/libxapp/xapp-util.h
lib/girepository-1.0/XApp-1.0.typelib
lib/gtk-3.0/modules/libxapp-gtk3-module.so
lib/libxapp.so
lib/libxapp.so.1
lib/libxapp.so.%%LIBVERSION%%
%%PYTHON_SITELIBDIR%%/gi/overrides/XApp.py
libdata/pkgconfig/xapp.pc
libexec/xapps/applet_constants.py
libexec/xapps/mate-xapp-status-applet.py
libexec/xapps/sn-watcher/xapp-sn-watcher
share/dbus-1/services/org.mate.panel.applet.MateXAppStatusAppletFactory.service
share/dbus-1/services/org.x.StatusNotifierWatcher.service
share/gir-1.0/XApp-1.0.gir
share/glade/catalogs/xapp-glade-catalog.xml
%%PORTDOCS%%share/gtk-doc/html/libxapp/XAppFavorites.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/XAppGtkWindow.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/XAppIconChooserButton.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/XAppIconChooserDialog.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/XAppKbdLayoutController.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/XAppMonitorBlanker.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/XAppObjectManagerClient.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/XAppPreferencesWindow.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/XAppStackSidebar.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/XAppStatusIcon.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/XAppStatusIconMonitor.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/annotation-glossary.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/api-index-full.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/ch01.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/deprecated-api-index.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/home.png
%%PORTDOCS%%share/gtk-doc/html/libxapp/index.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/left-insensitive.png
%%PORTDOCS%%share/gtk-doc/html/libxapp/left.png
%%PORTDOCS%%share/gtk-doc/html/libxapp/libxapp.devhelp2
%%PORTDOCS%%share/gtk-doc/html/libxapp/libxapp-xapp-util.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/object-tree.html
%%PORTDOCS%%share/gtk-doc/html/libxapp/right-insensitive.png
%%PORTDOCS%%share/gtk-doc/html/libxapp/right.png
%%PORTDOCS%%share/gtk-doc/html/libxapp/style.css
%%PORTDOCS%%share/gtk-doc/html/libxapp/up-insensitive.png
%%PORTDOCS%%share/gtk-doc/html/libxapp/up.png
share/icons/hicolor/scalable/actions/add-files-to-archive-symbolic.svg
share/icons/hicolor/scalable/actions/category-search-symbolic.svg
share/icons/hicolor/scalable/actions/extract-archive-symbolic.svg
share/icons/hicolor/scalable/actions/media-mount-symbolic.svg
share/icons/hicolor/scalable/actions/view-bottom-pane-symbolic.svg
share/icons/hicolor/scalable/actions/view-left-pane-symbolic-rtl.svg
share/icons/hicolor/scalable/actions/view-left-pane-symbolic.svg
share/icons/hicolor/scalable/actions/view-right-pane-symbolic-rtl.svg
share/icons/hicolor/scalable/actions/view-right-pane-symbolic.svg
share/icons/hicolor/scalable/actions/view-top-pane-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-annotations-squiggly-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-annotations-text-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-diaporama-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-favorite-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-format-text-highlight-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-go-history-next-symbolic-rtl.svg
share/icons/hicolor/scalable/actions/xapp-go-history-next-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-go-history-previous-symbolic-rtl.svg
share/icons/hicolor/scalable/actions/xapp-go-history-previous-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-pin-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-search-entire-word-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-search-wrap-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-text-case-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-unfavorite-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-unpin-symbolic.svg
share/icons/hicolor/scalable/actions/xapp-use-regex-symbolic.svg
share/icons/hicolor/scalable/apps/xapp-favorites-app-symbolic.svg
share/icons/hicolor/scalable/apps/xapp-favorites-app.svg
share/icons/hicolor/scalable/categories/xapp-prefs-behavior-symbolic.svg
share/icons/hicolor/scalable/categories/xapp-prefs-display-symbolic.svg
share/icons/hicolor/scalable/categories/xapp-prefs-plugins-symbolic.svg
share/icons/hicolor/scalable/categories/xapp-prefs-preview-symbolic.svg
share/icons/hicolor/scalable/categories/xapp-prefs-toolbar-symbolic.svg
share/icons/hicolor/scalable/emblems/emblem-xapp-favorite.svg
share/icons/hicolor/scalable/places/xapp-user-favorites-symbolic.svg
share/icons/hicolor/scalable/places/xapp-user-favorites.svg
share/locale/af/LC_MESSAGES/xapp.mo
share/locale/am/LC_MESSAGES/xapp.mo
share/locale/ar/LC_MESSAGES/xapp.mo
share/locale/az/LC_MESSAGES/xapp.mo
share/locale/ber/LC_MESSAGES/xapp.mo
share/locale/bg/LC_MESSAGES/xapp.mo
share/locale/bn/LC_MESSAGES/xapp.mo
share/locale/ca/LC_MESSAGES/xapp.mo
share/locale/ca@valencia/LC_MESSAGES/xapp.mo
share/locale/cs/LC_MESSAGES/xapp.mo
share/locale/cy/LC_MESSAGES/xapp.mo
share/locale/da/LC_MESSAGES/xapp.mo
share/locale/de/LC_MESSAGES/xapp.mo
share/locale/el/LC_MESSAGES/xapp.mo
share/locale/en_GB/LC_MESSAGES/xapp.mo
share/locale/eo/LC_MESSAGES/xapp.mo
share/locale/es/LC_MESSAGES/xapp.mo
share/locale/et/LC_MESSAGES/xapp.mo
share/locale/eu/LC_MESSAGES/xapp.mo
share/locale/fi/LC_MESSAGES/xapp.mo
share/locale/fr/LC_MESSAGES/xapp.mo
share/locale/fr_CA/LC_MESSAGES/xapp.mo
share/locale/he/LC_MESSAGES/xapp.mo
share/locale/hi/LC_MESSAGES/xapp.mo
share/locale/hr/LC_MESSAGES/xapp.mo
share/locale/hu/LC_MESSAGES/xapp.mo
share/locale/ia/LC_MESSAGES/xapp.mo
share/locale/id/LC_MESSAGES/xapp.mo
share/locale/ie/LC_MESSAGES/xapp.mo
share/locale/is/LC_MESSAGES/xapp.mo
share/locale/it/LC_MESSAGES/xapp.mo
share/locale/ja/LC_MESSAGES/xapp.mo
share/locale/ka/LC_MESSAGES/xapp.mo
share/locale/kab/LC_MESSAGES/xapp.mo
share/locale/kk/LC_MESSAGES/xapp.mo
share/locale/ko/LC_MESSAGES/xapp.mo
share/locale/ku/LC_MESSAGES/xapp.mo
share/locale/lt/LC_MESSAGES/xapp.mo
share/locale/mai/LC_MESSAGES/xapp.mo
share/locale/ms/LC_MESSAGES/xapp.mo
share/locale/nl/LC_MESSAGES/xapp.mo
share/locale/oc/LC_MESSAGES/xapp.mo
share/locale/pl/LC_MESSAGES/xapp.mo
share/locale/pt/LC_MESSAGES/xapp.mo
share/locale/pt_BR/LC_MESSAGES/xapp.mo
share/locale/ro/LC_MESSAGES/xapp.mo
share/locale/ru/LC_MESSAGES/xapp.mo
share/locale/sk/LC_MESSAGES/xapp.mo
share/locale/sl/LC_MESSAGES/xapp.mo
share/locale/sr/LC_MESSAGES/xapp.mo
share/locale/sr@latin/LC_MESSAGES/xapp.mo
share/locale/sv/LC_MESSAGES/xapp.mo
share/locale/tr/LC_MESSAGES/xapp.mo
share/locale/uk/LC_MESSAGES/xapp.mo
share/locale/uz/LC_MESSAGES/xapp.mo
share/locale/zgh/LC_MESSAGES/xapp.mo
share/locale/zh_CN/LC_MESSAGES/xapp.mo
share/locale/zh_HK/LC_MESSAGES/xapp.mo
share/locale/zh_TW/LC_MESSAGES/xapp.mo
share/mate-panel/applets/org.x.MateXAppStatusApplet.mate-panel-applet
share/vala/vapi/xapp.deps
share/vala/vapi/xapp.vapi
