# www/opengist/pkg-plist
bin/opengist
%%ETCDIR%%/config.yml.sample
etc/rc.d/opengist
%%WWWDIR%%/assets/KaTeX_AMS-Regular-BQhdFMY1.woff2
%%WWWDIR%%/assets/KaTeX_AMS-Regular-DMm9YOAa.woff
%%WWWDIR%%/assets/KaTeX_AMS-Regular-DRggAlZN.ttf
%%WWWDIR%%/assets/KaTeX_Caligraphic-Bold-ATXxdsX0.ttf
%%WWWDIR%%/assets/KaTeX_Caligraphic-Bold-BEiXGLvX.woff
%%WWWDIR%%/assets/KaTeX_Caligraphic-Bold-Dq_IR9rO.woff2
%%WWWDIR%%/assets/KaTeX_Caligraphic-Regular-CTRA-rTL.woff
%%WWWDIR%%/assets/KaTeX_Caligraphic-Regular-Di6jR-x-.woff2
%%WWWDIR%%/assets/KaTeX_Caligraphic-Regular-wX97UBjC.ttf
%%WWWDIR%%/assets/KaTeX_Fraktur-Bold-BdnERNNW.ttf
%%WWWDIR%%/assets/KaTeX_Fraktur-Bold-BsDP51OF.woff
%%WWWDIR%%/assets/KaTeX_Fraktur-Bold-CL6g_b3V.woff2
%%WWWDIR%%/assets/KaTeX_Fraktur-Regular-CB_wures.ttf
%%WWWDIR%%/assets/KaTeX_Fraktur-Regular-CTYiF6lA.woff2
%%WWWDIR%%/assets/KaTeX_Fraktur-Regular-Dxdc4cR9.woff
%%WWWDIR%%/assets/KaTeX_Main-Bold-Cx986IdX.woff2
%%WWWDIR%%/assets/KaTeX_Main-Bold-Jm3AIy58.woff
%%WWWDIR%%/assets/KaTeX_Main-Bold-waoOVXN0.ttf
%%WWWDIR%%/assets/KaTeX_Main-BoldItalic-DxDJ3AOS.woff2
%%WWWDIR%%/assets/KaTeX_Main-BoldItalic-DzxPMmG6.ttf
%%WWWDIR%%/assets/KaTeX_Main-BoldItalic-SpSLRI95.woff
%%WWWDIR%%/assets/KaTeX_Main-Italic-3WenGoN9.ttf
%%WWWDIR%%/assets/KaTeX_Main-Italic-BMLOBm91.woff
%%WWWDIR%%/assets/KaTeX_Main-Italic-NWA7e6Wa.woff2
%%WWWDIR%%/assets/KaTeX_Main-Regular-B22Nviop.woff2
%%WWWDIR%%/assets/KaTeX_Main-Regular-Dr94JaBh.woff
%%WWWDIR%%/assets/KaTeX_Main-Regular-ypZvNtVU.ttf
%%WWWDIR%%/assets/KaTeX_Math-BoldItalic-B3XSjfu4.ttf
%%WWWDIR%%/assets/KaTeX_Math-BoldItalic-CZnvNsCZ.woff2
%%WWWDIR%%/assets/KaTeX_Math-BoldItalic-iY-2wyZ7.woff
%%WWWDIR%%/assets/KaTeX_Math-Italic-DA0__PXp.woff
%%WWWDIR%%/assets/KaTeX_Math-Italic-flOr_0UB.ttf
%%WWWDIR%%/assets/KaTeX_Math-Italic-t53AETM-.woff2
%%WWWDIR%%/assets/KaTeX_SansSerif-Bold-CFMepnvq.ttf
%%WWWDIR%%/assets/KaTeX_SansSerif-Bold-D1sUS0GD.woff2
%%WWWDIR%%/assets/KaTeX_SansSerif-Bold-DbIhKOiC.woff
%%WWWDIR%%/assets/KaTeX_SansSerif-Italic-C3H0VqGB.woff2
%%WWWDIR%%/assets/KaTeX_SansSerif-Italic-DN2j7dab.woff
%%WWWDIR%%/assets/KaTeX_SansSerif-Italic-YYjJ1zSn.ttf
%%WWWDIR%%/assets/KaTeX_SansSerif-Regular-BNo7hRIc.ttf
%%WWWDIR%%/assets/KaTeX_SansSerif-Regular-CS6fqUqJ.woff
%%WWWDIR%%/assets/KaTeX_SansSerif-Regular-DDBCnlJ7.woff2
%%WWWDIR%%/assets/KaTeX_Script-Regular-C5JkGWo-.ttf
%%WWWDIR%%/assets/KaTeX_Script-Regular-D3wIWfF6.woff2
%%WWWDIR%%/assets/KaTeX_Script-Regular-D5yQViql.woff
%%WWWDIR%%/assets/KaTeX_Size1-Regular-C195tn64.woff
%%WWWDIR%%/assets/KaTeX_Size1-Regular-Dbsnue_I.ttf
%%WWWDIR%%/assets/KaTeX_Size1-Regular-mCD8mA8B.woff2
%%WWWDIR%%/assets/KaTeX_Size2-Regular-B7gKUWhC.ttf
%%WWWDIR%%/assets/KaTeX_Size2-Regular-Dy4dx90m.woff2
%%WWWDIR%%/assets/KaTeX_Size2-Regular-oD1tc_U0.woff
%%WWWDIR%%/assets/KaTeX_Size3-Regular-CTq5MqoE.woff
%%WWWDIR%%/assets/KaTeX_Size3-Regular-DgpXs0kz.ttf
%%WWWDIR%%/assets/KaTeX_Size3-Regular-gV2CO0n9.woff2
%%WWWDIR%%/assets/KaTeX_Size4-Regular-BF-4gkZK.woff
%%WWWDIR%%/assets/KaTeX_Size4-Regular-DWFBv043.ttf
%%WWWDIR%%/assets/KaTeX_Size4-Regular-Dl5lxZxV.woff2
%%WWWDIR%%/assets/KaTeX_Typewriter-Regular-C0xS9mPB.woff
%%WWWDIR%%/assets/KaTeX_Typewriter-Regular-CO6r4hn1.woff2
%%WWWDIR%%/assets/KaTeX_Typewriter-Regular-D3Ib7_Hf.ttf
%%WWWDIR%%/assets/_commonjsHelpers-CqkleIqs.js
%%WWWDIR%%/assets/admin-zSosdAUq.js
%%WWWDIR%%/assets/auto-AI_rlwP9.js
%%WWWDIR%%/assets/auto-Cgp9ogWp.css
%%WWWDIR%%/assets/dark-lPV00ZsJ.js
%%WWWDIR%%/assets/dark-yXNGFPPS.css
%%WWWDIR%%/assets/editor-DVd6Cy7T.js
%%WWWDIR%%/assets/embed-BxqsL3ZH.js
%%WWWDIR%%/assets/embed-CuaRj9uV.css
%%WWWDIR%%/assets/favicon-32-BHd4osPo.png
%%WWWDIR%%/assets/gist-8fqheCgo.js
%%WWWDIR%%/assets/light-B0neenoN.js
%%WWWDIR%%/assets/light-BHyWC2Xy.css
%%WWWDIR%%/assets/main-BOjJ7ls3.js
%%WWWDIR%%/assets/main-BeRbbMrT.css
%%WWWDIR%%/assets/opengist-DI_FK2dF.svg
%%WWWDIR%%/assets/style_preferences-Ck4u9x6z.js
%%WWWDIR%%/assets/webauthn-Bx0wA4K6.js
@dir(opengist,opengist,750) /var/db/opengist
@dir(opengist,opengist,750) /var/run/opengist
