# net/openntpd/pkg-plist
@comment $FreeBSD$
etc/rc.d/openntpd.sh
sbin/ntpd
@unexec if cmp -s %D/etc/ntpd.conf %D/%%EXAMPLESDIR%%/ntpd.conf; then rm -f %D/etc/ntpd.conf; fi
%%EXAMPLESDIR%%/ntpd.conf
@exec if [ ! -f %D/etc/ntpd.conf ]; then cp -p %D/%F %D/etc; fi
@dirrm %%EXAMPLESDIR%%
