# sysutils/pcbsd-appweb/pkg-plist
bin/appcafe-setpass
etc/appcafe.conf.dist
etc/sudoers.d/appcafe-sudoers
share/appcafe/dispatcher
share/appcafe/dispatcher-auth
share/appcafe/dispatcher-localauth
share/appcafe/fastcgi_params
share/appcafe/mime.types
share/appcafe/nginx.conf
share/appcafe/nginx.conf.ssl
share/appcafe/nginx.conf.temp
share/appcafe/html/index.php
share/appcafe/html/css/_reset.css
share/appcafe/html/css/cssmenu.css
share/appcafe/html/css/main.css
share/appcafe/html/css/popbox.css
share/appcafe/html/css/tabs.css
share/appcafe/html/css/lightbox.css
share/appcafe/html/images/add.png
share/appcafe/html/images/appcafe.png
share/appcafe/html/images/application-exit.png
share/appcafe/html/images/categories.png
share/appcafe/html/images/close.png
share/appcafe/html/images/configure.png
share/appcafe/html/images/cssbg.png
share/appcafe/html/images/dialog-ok.png
share/appcafe/html/images/external-link.png
share/appcafe/html/images/info-tips.png
share/appcafe/html/images/install.png
share/appcafe/html/images/jail.png
share/appcafe/html/images/loading.gif
share/appcafe/html/images/logout.png
share/appcafe/html/images/next.png
share/appcafe/html/images/notifier.png
share/appcafe/html/images/noupdates.png
share/appcafe/html/images/pbiicon.php
share/appcafe/html/images/prev.png
share/appcafe/html/images/rating-0.png
share/appcafe/html/images/rating-1.png
share/appcafe/html/images/rating-2.png
share/appcafe/html/images/rating-3.png
share/appcafe/html/images/rating-4.png
share/appcafe/html/images/rating-5.png
share/appcafe/html/images/remove.png
share/appcafe/html/images/restart.png
share/appcafe/html/images/search.png
share/appcafe/html/images/start.png
share/appcafe/html/images/switch.png
share/appcafe/html/images/system.png
share/appcafe/html/images/warden.png
share/appcafe/html/images/warning.png
share/appcafe/html/images/working.gif
share/appcafe/html/include/Mobile_Detect.php
share/appcafe/html/include/header.php
share/appcafe/html/include/functions.php
share/appcafe/html/include/footer.php
share/appcafe/html/include/globals.php
share/appcafe/html/include/leftmenu.php
share/appcafe/html/include/login.php
share/appcafe/html/include/nav-computer.php
share/appcafe/html/include/nav-phone.php
share/appcafe/html/pages/appcafe.php
share/appcafe/html/pages/appcafe-browse.php
share/appcafe/html/pages/appcafe-search.php
share/appcafe/html/pages/appinfo.php
share/appcafe/html/pages/dispatcher.php
share/appcafe/html/pages/error.php
share/appcafe/html/pages/jails.php
share/appcafe/html/pages/jailcreate.php
share/appcafe/html/pages/jailinfo.php
share/appcafe/html/pages/main.php
share/appcafe/html/pages/notifier.php
share/appcafe/html/pages/sysapp.php
share/appcafe/html/js/appcafe.js
share/appcafe/html/js/cssmenu.js
share/appcafe/html/js/doubletaptogo.js
share/appcafe/html/js/jquery.ba-hashchange.min.js
share/appcafe/html/js/jquery.easytabs.min.js
share/appcafe/html/js/jquery.jscroll.min.js
share/appcafe/html/js/jquery.timer.js
share/appcafe/html/js/jquery.min.js
share/appcafe/html/js/lightbox.min.js
share/appcafe/html/js/main.js
share/appcafe/html/js/popbox.min.js
etc/rc.d/appcafe
