# net/rwhoisd/files/rwhoisd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: rwhoisd
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# rwhoisd_enable (bool):   Set to NO by default.
#               Set it to YES to enable rwhoisd.
# rwhoisd_config (path):   Set to %%PREFIX%%/etc/rwhoisd/rwhoisd.conf
#               by default.
#

. /etc/rc.subr

name=rwhoisd
rcvar=rwhoisd_enable
load_rc_config $name

: ${rwhoisd_enable:=NO}
: ${rwhoisd_config="%%PREFIX%%/etc/rwhoisd/rwhoisd.conf"}
: ${rwhoisd_pidfile="/var/run/${name}/${name}.pid"}

command=%%PREFIX%%/sbin/${name}
pidfile=${rwhoisd_pidfile}

command_args="-d -c $rwhoisd_config"

run_rc_command "$1"
