# net-mgmt/netams/pkg-plist
etc/netams.cfg.sample
etc/rc.d/netams-startup.sh.sample
libexec/netams
sbin/flowprobe
sbin/ipfw2netflow
bin/netamsctl
%%DATADIR%%/ru-networks.txt.sample
%%DATADIR%%/.netamsctl.rc
%%DATADIR%%/cgi-bin/.htaccess
%%DATADIR%%/cgi-bin/config.cgi
%%DATADIR%%/cgi-bin/netams_api.pl
%%DATADIR%%/cgi-bin/russian.res
%%DATADIR%%/cgi-bin/weblogin.tem
%%DATADIR%%/cgi-bin/netams_example.cgi
%%DATADIR%%/cgi-bin/usertool.cgi
%%DATADIR%%/cgi-bin/admintool.cgi
%%DATADIR%%/cgi-bin/login.cgi
%%DATADIR%%/cgi-bin/netams_graph.cgi
%%DATADIR%%/cgi-bin/weblogin.cgi
%%DATADIR%%/cgi-bin/admin/.htaccess
%%DATADIR%%/cgi-bin/admin/index.cgi
%%DATADIR%%/cgi-bin/admin/policy.cgi
%%DATADIR%%/cgi-bin/admin/showusercard.cgi
%%DATADIR%%/cgi-bin/admin/user.cgi
%%DATADIR%%/cgi-bin/admin/access.cgi
%%DATADIR%%/cgi-bin/admin/login.cgi
%%DATADIR%%/cgi-bin/admin/quota.cgi
%%DATADIR%%/cgi-bin/admin/showusercard.tmpl
%%DATADIR%%/cgi-bin/admin/view.cgi
%%DATADIR%%/cgi-bin/admin/account.cgi
%%DATADIR%%/cgi-bin/admin/netams.cgi
%%DATADIR%%/cgi-bin/admin/radius.cgi
%%DATADIR%%/cgi-bin/admin/subplan.cgi
%%DATADIR%%/cgi-bin/admin/config.cgi
%%DATADIR%%/cgi-bin/admin/plan.cgi
%%DATADIR%%/cgi-bin/admin/russian.res
%%DATADIR%%/cgi-bin/admin/unit.cgi
%%DATADIR%%/cgi-bin/images/admintool-logo.gif
%%DATADIR%%/cgi-bin/images/logo-small.gif
%%DATADIR%%/cgi-bin/images/logo.gif
%%DATADIR%%/cgi-bin/images/showtable-logo.gif
@dirrm %%DATADIR%%/cgi-bin/images
@dirrm %%DATADIR%%/cgi-bin/admin
@dirrm %%DATADIR%%/cgi-bin
@dirrm %%DATADIR%%
