# net/pimd/pkg-plist
@sample etc/pimd.conf.sample
etc/rc.d/pimd
sbin/pimd
man/man8/pimd.8.gz
