# security/expiretable/files/openrc-expiretable.in
#!/sbin/openrc-run

# /etc/rc.conf tunables. 
# - expiretable_flags	specify your own flags to expiretable

name="expiretable"
command="/usr/local/sbin/expiretable"

depend() {
		keyword -shutdown
}
