# dns/ez-ipupdate/pkg-plist
bin/ez-ipupdate
etc/ez-ipupdate.conf.sample
etc/rc.d/ez-ipupdate.sh
share/examples/ez-ipupdate/example-dhs.conf
share/examples/ez-ipupdate/example-dyndns.conf
share/examples/ez-ipupdate/example-dyns.conf
share/examples/ez-ipupdate/example-easydns.conf
share/examples/ez-ipupdate/example-gnudip.conf
share/examples/ez-ipupdate/example-heipv6tb.conf
share/examples/ez-ipupdate/example-justlinux.conf
share/examples/ez-ipupdate/example-ods.conf
share/examples/ez-ipupdate/example-pgpow.conf
share/examples/ez-ipupdate/example-tzo.conf
@dirrm share/examples/ez-ipupdate
