# graphics/linux-rl9-dri/pkg-plist.amd64
usr/lib/.build-id/04/32745fc446b7656ad2224f1706b6f4aa902b73
usr/lib/.build-id/04/32745fc446b7656ad2224f1706b6f4aa902b73.1
usr/lib/.build-id/04/32745fc446b7656ad2224f1706b6f4aa902b73.2
usr/lib/.build-id/17/afed4f3c0dc1eaa4d2ac486b0c2254bb5ff9c8
usr/lib/.build-id/2c/d6463aabf53225f2066304ca273ac4064f9715
usr/lib/.build-id/34/882966bada406a811cdc5ddf3ab7388faf2de9
usr/lib/.build-id/4c/487cdf156eb3813f9252cce2caae3991716795
usr/lib/.build-id/4c/487cdf156eb3813f9252cce2caae3991716795.1
usr/lib/.build-id/4c/487cdf156eb3813f9252cce2caae3991716795.2
usr/lib/.build-id/4c/487cdf156eb3813f9252cce2caae3991716795.3
usr/lib/.build-id/4c/487cdf156eb3813f9252cce2caae3991716795.4
usr/lib/.build-id/4c/487cdf156eb3813f9252cce2caae3991716795.5
usr/lib/.build-id/4c/487cdf156eb3813f9252cce2caae3991716795.6
usr/lib/.build-id/4c/487cdf156eb3813f9252cce2caae3991716795.7
usr/lib/.build-id/4c/487cdf156eb3813f9252cce2caae3991716795.8
usr/lib/.build-id/4d/dcf3cb29b4bfefcaa49e996e57a621a809f3ff
usr/lib/.build-id/4e/c92eb39b18f69ad6eeb091281bf93a80908064
usr/lib/.build-id/51/de0e8f1ef95feafc18270ec8af3a0e3844e87b
usr/lib/.build-id/82/7d94e2b6bd9881a72fe41df69ef6e587a518dd
usr/lib/.build-id/9b/e22acbc9fcfc465f068e5f9eeb4e2b2156ed0d
usr/lib/.build-id/b2/28b422da64ec6945810689286963494c7ad41a
usr/lib/.build-id/cd/d2f7afd304586475b15ad7e81c08bfe01ead33
usr/lib/.build-id/f1/c0b3d1460a35848cb14a9e51b2bde3351cd132
usr/lib/.build-id/f7/e9286e383fbe12823da2e0164551662849e28b
usr/lib/libGLX_system.so.0
usr/lib64/dri/crocus_dri.so
usr/lib64/dri/iris_dri.so
usr/lib64/dri/kms_swrast_dri.so
usr/lib64/dri/nouveau_dri.so
usr/lib64/dri/nouveau_drv_video.so
usr/lib64/dri/radeonsi_dri.so
usr/lib64/dri/radeonsi_drv_video.so
usr/lib64/dri/swrast_dri.so
usr/lib64/dri/virtio_gpu_dri.so
usr/lib64/dri/virtio_gpu_drv_video.so
usr/lib64/dri/vmwgfx_dri.so
usr/lib64/dri/zink_dri.so
usr/lib64/libEGL_mesa.so.0
usr/lib64/libEGL_mesa.so.0.0.0
usr/lib64/libGLU.so.1
usr/lib64/libGLU.so.1.3.1
usr/lib64/libGLX_mesa.so.0
usr/lib64/libGLX_mesa.so.0.0.0
usr/lib64/libGLX_system.so.0
usr/lib64/libVkLayer_MESA_device_select.so
usr/lib64/libgbm.so.1
usr/lib64/libgbm.so.1.0.0
usr/lib64/libglapi.so.0
usr/lib64/libglapi.so.0.0.0
usr/lib64/libglut.so.3
usr/lib64/libglut.so.3.11.0
usr/lib64/libvulkan_intel.so
usr/lib64/libvulkan_intel_hasvk.so
usr/lib64/libvulkan_lvp.so
usr/lib64/libvulkan_radeon.so
usr/lib64/libxatracker.so.2
usr/lib64/libxatracker.so.2.5.0
%%PORTDOCS%%usr/share/doc/freeglut/AUTHORS
%%PORTDOCS%%usr/share/doc/freeglut/COPYING
%%PORTDOCS%%usr/share/doc/freeglut/ChangeLog
%%PORTDOCS%%usr/share/doc/freeglut/README
%%PORTDOCS%%usr/share/doc/freeglut/download.html
%%PORTDOCS%%usr/share/doc/freeglut/freeglut.html
%%PORTDOCS%%usr/share/doc/freeglut/freeglut_logo.png
%%PORTDOCS%%usr/share/doc/freeglut/freeglut_user_interface.html
%%PORTDOCS%%usr/share/doc/freeglut/index.html
%%PORTDOCS%%usr/share/doc/freeglut/ogl_sm.png
%%PORTDOCS%%usr/share/doc/freeglut/progress.html
%%PORTDOCS%%usr/share/doc/freeglut/structure.html
%%DATADIR%%rc.d/00-mesa-defaults.conf
%%DATADIR%%rc.d/00-radv-defaults.conf
usr/share/glvnd/egl_vendor.d/50_mesa.json
usr/share/vulkan/icd.d/intel_hasvk_icd.x86_64.json
usr/share/vulkan/icd.d/intel_icd.x86_64.json
usr/share/vulkan/icd.d/lvp_icd.x86_64.json
usr/share/vulkan/icd.d/radeon_icd.x86_64.json
usr/share/vulkan/implicit_layer.d/VkLayer_MESA_device_select.json
