# graphics/linux-rl9-dri/pkg-plist.aarch64
usr/lib/.build-id/10/5d263ef808549ba0ad9fb95cca20232cf8d53f
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.1
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.10
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.11
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.12
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.13
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.14
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.15
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.16
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.17
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.18
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.19
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.2
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.20
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.21
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.22
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.23
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.24
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.25
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.26
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.27
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.28
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.29
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.3
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.30
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.31
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.32
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.33
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.34
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.4
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.5
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.6
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.7
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.8
usr/lib/.build-id/1f/8cfbfed314476d0fdcaed37691eaeb7e14a1a8.9
usr/lib/.build-id/47/857a292c3d46db9d66b139e2a2fb9c90894c39
usr/lib/.build-id/53/e67200830cee99776abe0709217d7ced3e9cdd
usr/lib/.build-id/a0/dd4e3b2b5b22d09fe5c2b4c78fe50cad165d83
usr/lib/.build-id/b2/38fb37e9652a9e4e4e19450df44968f259d26a
usr/lib/.build-id/b2/38fb37e9652a9e4e4e19450df44968f259d26a.1
usr/lib/.build-id/b2/38fb37e9652a9e4e4e19450df44968f259d26a.2
usr/lib/.build-id/ec/f38088658bb9156dc989a7ea9bdb6efb90d1c0
usr/lib/.build-id/fc/f7b0d5b101f4bcc3ecaedb8887d99e81df0371
usr/lib64/dri/armada-drm_dri.so
usr/lib64/dri/exynos_dri.so
usr/lib64/dri/hx8357d_dri.so
usr/lib64/dri/ili9225_dri.so
usr/lib64/dri/ili9341_dri.so
usr/lib64/dri/imx-dcss_dri.so
usr/lib64/dri/imx-drm_dri.so
usr/lib64/dri/imx-lcdif_dri.so
usr/lib64/dri/ingenic-drm_dri.so
usr/lib64/dri/kgsl_dri.so
usr/lib64/dri/kirin_dri.so
usr/lib64/dri/kms_swrast_dri.so
usr/lib64/dri/komeda_dri.so
usr/lib64/dri/mali-dp_dri.so
usr/lib64/dri/mcde_dri.so
usr/lib64/dri/mediatek_dri.so
usr/lib64/dri/meson_dri.so
usr/lib64/dri/mi0283qt_dri.so
usr/lib64/dri/msm_dri.so
usr/lib64/dri/mxsfb-drm_dri.so
usr/lib64/dri/nouveau_dri.so
usr/lib64/dri/nouveau_drv_video.so
usr/lib64/dri/panfrost_dri.so
usr/lib64/dri/pl111_dri.so
usr/lib64/dri/radeonsi_dri.so
usr/lib64/dri/radeonsi_drv_video.so
usr/lib64/dri/rcar-du_dri.so
usr/lib64/dri/repaper_dri.so
usr/lib64/dri/rockchip_dri.so
usr/lib64/dri/st7586_dri.so
usr/lib64/dri/st7735r_dri.so
usr/lib64/dri/stm_dri.so
usr/lib64/dri/sun4i-drm_dri.so
usr/lib64/dri/swrast_dri.so
usr/lib64/dri/tegra_dri.so
usr/lib64/dri/virtio_gpu_dri.so
usr/lib64/dri/virtio_gpu_drv_video.so
usr/lib64/dri/zink_dri.so
usr/lib64/libEGL_mesa.so.0
usr/lib64/libEGL_mesa.so.0.0.0
usr/lib64/libGLU.so.1
usr/lib64/libGLU.so.1.3.1
usr/lib64/libGLX_mesa.so.0
usr/lib64/libGLX_mesa.so.0.0.0
usr/lib64/libGLX_system.so.0
usr/lib64/libgbm.so.1
usr/lib64/libgbm.so.1.0.0
usr/lib64/libglapi.so.0
usr/lib64/libglapi.so.0.0.0
usr/lib64/libglut.so.3
usr/lib64/libglut.so.3.11.0
%%PORTDOCS%%usr/share/doc/freeglut/AUTHORS
%%PORTDOCS%%usr/share/doc/freeglut/COPYING
%%PORTDOCS%%usr/share/doc/freeglut/ChangeLog
%%PORTDOCS%%usr/share/doc/freeglut/README
%%PORTDOCS%%usr/share/doc/freeglut/download.html
%%PORTDOCS%%usr/share/doc/freeglut/freeglut.html
%%PORTDOCS%%usr/share/doc/freeglut/freeglut_logo.png
%%PORTDOCS%%usr/share/doc/freeglut/freeglut_user_interface.html
%%PORTDOCS%%usr/share/doc/freeglut/index.html
%%PORTDOCS%%usr/share/doc/freeglut/ogl_sm.png
%%PORTDOCS%%usr/share/doc/freeglut/progress.html
%%PORTDOCS%%usr/share/doc/freeglut/structure.html
%%DATADIR%%rc.d/00-mesa-defaults.conf
usr/share/glvnd/egl_vendor.d/50_mesa.json
