# graphics/linux-rl9-dri/pkg-plist.amd64
usr/lib/.build-id/00/1ef622449b65e9b9efcd3486ad41ee7b3029ed
usr/lib/.build-id/0f/fb1d7f4aefa7d96a42090a6299bbe574699990
usr/lib/.build-id/27/db423e60c886b9a4ee72e2f686182deb0a83e7
usr/lib/.build-id/27/db423e60c886b9a4ee72e2f686182deb0a83e7.1
usr/lib/.build-id/27/db423e60c886b9a4ee72e2f686182deb0a83e7.2
usr/lib/.build-id/27/db423e60c886b9a4ee72e2f686182deb0a83e7.3
usr/lib/.build-id/27/db423e60c886b9a4ee72e2f686182deb0a83e7.4
usr/lib/.build-id/27/db423e60c886b9a4ee72e2f686182deb0a83e7.5
usr/lib/.build-id/27/db423e60c886b9a4ee72e2f686182deb0a83e7.6
usr/lib/.build-id/27/db423e60c886b9a4ee72e2f686182deb0a83e7.7
usr/lib/.build-id/27/db423e60c886b9a4ee72e2f686182deb0a83e7.8
usr/lib/.build-id/33/c0d3af0221ea652b57ed45447d72b4f3f9e13a
usr/lib/.build-id/3f/c25f71e1b4794a4bf999bc021baa739223e6ac
usr/lib/.build-id/50/8524465b49e38e02f673281244d9ef0fe668a9
usr/lib/.build-id/50/8524465b49e38e02f673281244d9ef0fe668a9.1
usr/lib/.build-id/50/8524465b49e38e02f673281244d9ef0fe668a9.2
usr/lib/.build-id/50/945cbbf1aa6d5d51132617a4a9ca18c0910c68
usr/lib/.build-id/51/de0e8f1ef95feafc18270ec8af3a0e3844e87b
usr/lib/.build-id/56/5591a0e8645fe70331b9a92bb480a8b0ce0aed
usr/lib/.build-id/56/ecf7f7c2f46c6e9471904e277c34a812bebf17
usr/lib/.build-id/73/ae5df50221ca4f575499af452a45d7847aa7de
usr/lib/.build-id/c4/7ef3cb628c8cba6eb123bc7ac6cce8115af612
usr/lib/.build-id/f0/990aa64df2c0a21ace4adb671b17e04472d890
usr/lib/.build-id/f1/c0b3d1460a35848cb14a9e51b2bde3351cd132
usr/lib/libGLX_system.so.0
usr/lib64/dri/crocus_dri.so
usr/lib64/dri/iris_dri.so
usr/lib64/dri/kms_swrast_dri.so
usr/lib64/dri/nouveau_dri.so
usr/lib64/dri/nouveau_drv_video.so
usr/lib64/dri/radeonsi_dri.so
usr/lib64/dri/radeonsi_drv_video.so
usr/lib64/dri/swrast_dri.so
usr/lib64/dri/virtio_gpu_dri.so
usr/lib64/dri/virtio_gpu_drv_video.so
usr/lib64/dri/vmwgfx_dri.so
usr/lib64/dri/zink_dri.so
usr/lib64/libEGL_mesa.so.0
usr/lib64/libEGL_mesa.so.0.0.0
usr/lib64/libGLU.so.1
usr/lib64/libGLU.so.1.3.1
usr/lib64/libGLX_mesa.so.0
usr/lib64/libGLX_mesa.so.0.0.0
usr/lib64/libGLX_system.so.0
usr/lib64/libVkLayer_MESA_device_select.so
usr/lib64/libgbm.so.1
usr/lib64/libgbm.so.1.0.0
usr/lib64/libglapi.so.0
usr/lib64/libglapi.so.0.0.0
usr/lib64/libglut.so.3
usr/lib64/libglut.so.3.11.0
usr/lib64/libvulkan_intel.so
usr/lib64/libvulkan_intel_hasvk.so
usr/lib64/libvulkan_lvp.so
usr/lib64/libvulkan_radeon.so
usr/lib64/libxatracker.so.2
usr/lib64/libxatracker.so.2.5.0
%%PORTDOCS%%usr/share/doc/freeglut/AUTHORS
%%PORTDOCS%%usr/share/doc/freeglut/COPYING
%%PORTDOCS%%usr/share/doc/freeglut/ChangeLog
%%PORTDOCS%%usr/share/doc/freeglut/README
%%PORTDOCS%%usr/share/doc/freeglut/download.html
%%PORTDOCS%%usr/share/doc/freeglut/freeglut.html
%%PORTDOCS%%usr/share/doc/freeglut/freeglut_logo.png
%%PORTDOCS%%usr/share/doc/freeglut/freeglut_user_interface.html
%%PORTDOCS%%usr/share/doc/freeglut/index.html
%%PORTDOCS%%usr/share/doc/freeglut/ogl_sm.png
%%PORTDOCS%%usr/share/doc/freeglut/progress.html
%%PORTDOCS%%usr/share/doc/freeglut/structure.html
%%DATADIR%%rc.d/00-mesa-defaults.conf
%%DATADIR%%rc.d/00-radv-defaults.conf
usr/share/glvnd/egl_vendor.d/50_mesa.json
usr/share/vulkan/icd.d/intel_hasvk_icd.x86_64.json
usr/share/vulkan/icd.d/intel_icd.x86_64.json
usr/share/vulkan/icd.d/lvp_icd.x86_64.json
usr/share/vulkan/icd.d/radeon_icd.x86_64.json
usr/share/vulkan/implicit_layer.d/VkLayer_MESA_device_select.json
