# graphics/linux-rl9-dri/pkg-plist.amd64
usr/lib/.build-id/0b/846527741b837edef4655045873fb175452d42
usr/lib/.build-id/0c/4d9154fce99746b6aaaffe40a5d2ca9288627d
usr/lib/.build-id/10/610b10ab8562a430d22d6e9fb359e778a7d12e
usr/lib/.build-id/26/3a8764ed3d987250c02d25eefa5aadb229000b
usr/lib/.build-id/30/e72f6af7da4d6da5d7c8ca394519b5228a1816
usr/lib/.build-id/35/e082d883adda79674d7c3b9bcc415d46f89851
usr/lib/.build-id/3a/36172eedd3c0a6f0f352a28ec817f4120f5fd2
usr/lib/.build-id/3a/36172eedd3c0a6f0f352a28ec817f4120f5fd2.1
usr/lib/.build-id/3a/36172eedd3c0a6f0f352a28ec817f4120f5fd2.2
usr/lib/.build-id/3a/36172eedd3c0a6f0f352a28ec817f4120f5fd2.3
usr/lib/.build-id/3a/36172eedd3c0a6f0f352a28ec817f4120f5fd2.4
usr/lib/.build-id/3a/36172eedd3c0a6f0f352a28ec817f4120f5fd2.5
usr/lib/.build-id/3a/36172eedd3c0a6f0f352a28ec817f4120f5fd2.6
usr/lib/.build-id/3a/36172eedd3c0a6f0f352a28ec817f4120f5fd2.7
usr/lib/.build-id/3a/36172eedd3c0a6f0f352a28ec817f4120f5fd2.8
usr/lib/.build-id/4d/6f0551a4d8b94f8b84e4af7b61994b9aa913ef
usr/lib/.build-id/51/de0e8f1ef95feafc18270ec8af3a0e3844e87b
usr/lib/.build-id/54/4a40184aee075712ae6bfdb53121abcf75a4a8
usr/lib/.build-id/54/58d9acf144487556ec41664d36a72838018e94
usr/lib/.build-id/61/7b6fe8008a91a422cb506ddd36818cc3de9de9
usr/lib/.build-id/7d/4777d32813159cb61c6a762c2121920d060268
usr/lib/.build-id/81/9b1dc0c85ff0564f791a377c0de04b6d7f8f81
usr/lib/.build-id/89/b25212de11e93f80ba05682a26de4419123ea6
usr/lib/.build-id/d7/ff0260c1c064b433a4f35b0588d461ea4a728d
usr/lib/.build-id/d8/7c33dd5908e4712c75bc8c8356e8ac678bb5d7
usr/lib/.build-id/d8/7c33dd5908e4712c75bc8c8356e8ac678bb5d7.1
usr/lib/.build-id/d8/7c33dd5908e4712c75bc8c8356e8ac678bb5d7.2
usr/lib/.build-id/d8/7c33dd5908e4712c75bc8c8356e8ac678bb5d7.3
usr/lib/.build-id/d8/7c33dd5908e4712c75bc8c8356e8ac678bb5d7.4
usr/lib/.build-id/d8/7c33dd5908e4712c75bc8c8356e8ac678bb5d7.5
usr/lib/.build-id/d8/7c33dd5908e4712c75bc8c8356e8ac678bb5d7.6
usr/lib/.build-id/d8/7c33dd5908e4712c75bc8c8356e8ac678bb5d7.7
usr/lib/.build-id/d8/7c33dd5908e4712c75bc8c8356e8ac678bb5d7.8
usr/lib/.build-id/dc/366f0d712458f8458f5d1feee2a1f1dfafe837
usr/lib/.build-id/de/c1c007bd7f5d80870dfaae3a50a8bf7e1c3a7e
usr/lib/.build-id/e6/024987291f86134749f3ed5e064e5b729942cf
usr/lib/dri/crocus_dri.so
usr/lib/dri/iris_dri.so
usr/lib/dri/kms_swrast_dri.so
usr/lib/dri/nouveau_dri.so
usr/lib/dri/radeonsi_dri.so
usr/lib/dri/swrast_dri.so
usr/lib/dri/virtio_gpu_dri.so
usr/lib/dri/vmwgfx_dri.so
usr/lib/dri/zink_dri.so
usr/lib/libEGL_mesa.so.0
usr/lib/libEGL_mesa.so.0.0.0
usr/lib/libGLU.so.1
usr/lib/libGLU.so.1.3.1
usr/lib/libGLX_mesa.so.0
usr/lib/libGLX_mesa.so.0.0.0
usr/lib/libGLX_system.so.0
usr/lib/libgbm.so.1
usr/lib/libgbm.so.1.0.0
usr/lib/libglapi.so.0
usr/lib/libglapi.so.0.0.0
usr/lib/libglut.so.3
usr/lib/libglut.so.3.11.0
usr/lib64/dri/crocus_dri.so
usr/lib64/dri/iris_dri.so
usr/lib64/dri/kms_swrast_dri.so
usr/lib64/dri/nouveau_dri.so
usr/lib64/dri/radeonsi_dri.so
usr/lib64/dri/swrast_dri.so
usr/lib64/dri/virtio_gpu_dri.so
usr/lib64/dri/vmwgfx_dri.so
usr/lib64/dri/zink_dri.so
usr/lib64/libEGL_mesa.so.0
usr/lib64/libEGL_mesa.so.0.0.0
usr/lib64/libGLU.so.1
usr/lib64/libGLU.so.1.3.1
usr/lib64/libGLX_mesa.so.0
usr/lib64/libGLX_mesa.so.0.0.0
usr/lib64/libGLX_system.so.0
usr/lib64/libVkLayer_MESA_device_select.so
usr/lib64/libgbm.so.1
usr/lib64/libgbm.so.1.0.0
usr/lib64/libglapi.so.0
usr/lib64/libglapi.so.0.0.0
usr/lib64/libglut.so.3
usr/lib64/libglut.so.3.11.0
usr/lib64/libvulkan_intel.so
usr/lib64/libvulkan_intel_hasvk.so
usr/lib64/libvulkan_lvp.so
usr/lib64/libvulkan_radeon.so
usr/lib64/libxatracker.so.2
usr/lib64/libxatracker.so.2.5.0
%%PORTDOCS%%usr/share/doc/freeglut/AUTHORS
%%PORTDOCS%%usr/share/doc/freeglut/COPYING
%%PORTDOCS%%usr/share/doc/freeglut/ChangeLog
%%PORTDOCS%%usr/share/doc/freeglut/README
%%PORTDOCS%%usr/share/doc/freeglut/download.html
%%PORTDOCS%%usr/share/doc/freeglut/freeglut.html
%%PORTDOCS%%usr/share/doc/freeglut/freeglut_logo.png
%%PORTDOCS%%usr/share/doc/freeglut/freeglut_user_interface.html
%%PORTDOCS%%usr/share/doc/freeglut/index.html
%%PORTDOCS%%usr/share/doc/freeglut/ogl_sm.png
%%PORTDOCS%%usr/share/doc/freeglut/progress.html
%%PORTDOCS%%usr/share/doc/freeglut/structure.html
%%DATADIR%%rc.d/00-mesa-defaults.conf
%%DATADIR%%rc.d/00-radv-defaults.conf
usr/share/glvnd/egl_vendor.d/50_mesa.json
usr/share/vulkan/icd.d/intel_hasvk_icd.x86_64.json
usr/share/vulkan/icd.d/intel_icd.x86_64.json
usr/share/vulkan/icd.d/lvp_icd.x86_64.json
usr/share/vulkan/icd.d/radeon_icd.x86_64.json
usr/share/vulkan/implicit_layer.d/VkLayer_MESA_device_select.json
