# graphics/linux-rl9-dri/pkg-plist.aarch64
usr/lib/.build-id/3a/876fff075f7e6e664b85bccfc4bd396db1b699
usr/lib/.build-id/40/f9445b135eefc58264383e035b8d554a2a677e
usr/lib/.build-id/7b/837bc45bf40edcf51b152a3ca1cf80770dce62
usr/lib/.build-id/84/ddcc9baa839267464a9d502367b458588b995e
usr/lib/.build-id/ab/52df3f33852384c6738a6e642e3f7b8dd5ccbb
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.1
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.2
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.3
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.4
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.5
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.6
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.7
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.8
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.9
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.10
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.11
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.12
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.13
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.14
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.15
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.16
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.17
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.18
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.19
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.20
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.21
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.22
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.23
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.24
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.25
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.26
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.27
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.28
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.29
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.30
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.31
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.32
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.33
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.34
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.35
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.36
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.37
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.38
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.39
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.40
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.41
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.42
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.43
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.44
usr/lib/.build-id/cc/35979fc234bdb78a680d3684627fcc50b62d7b.45
usr/lib/.build-id/ec/f38088658bb9156dc989a7ea9bdb6efb90d1c0
usr/lib64/dri/armada-drm_dri.so
usr/lib64/dri/exynos_dri.so
usr/lib64/dri/gm12u320_dri.so
usr/lib64/dri/hdlcd_dri.so
usr/lib64/dri/hx8357d_dri.so
usr/lib64/dri/ili9163_dri.so
usr/lib64/dri/ili9225_dri.so
usr/lib64/dri/ili9341_dri.so
usr/lib64/dri/ili9486_dri.so
usr/lib64/dri/imx-dcss_dri.so
usr/lib64/dri/imx-drm_dri.so
usr/lib64/dri/imx-lcdif_dri.so
usr/lib64/dri/ingenic-drm_dri.so
usr/lib64/dri/kgsl_dri.so
usr/lib64/dri/kirin_dri.so
usr/lib64/dri/kms_swrast_dri.so
usr/lib64/dri/komeda_dri.so
usr/lib64/dri/mali-dp_dri.so
usr/lib64/dri/mcde_dri.so
usr/lib64/dri/mediatek_dri.so
usr/lib64/dri/meson_dri.so
usr/lib64/dri/mi0283qt_dri.so
usr/lib64/dri/msm_dri.so
usr/lib64/dri/mxsfb-drm_dri.so
usr/lib64/dri/nouveau_dri.so
usr/lib64/dri/panel-mipi-dbi_dri.so
usr/lib64/dri/panfrost_dri.so
usr/lib64/dri/panthor_dri.so
usr/lib64/dri/pl111_dri.so
usr/lib64/dri/radeonsi_dri.so
usr/lib64/dri/rcar-du_dri.so
usr/lib64/dri/repaper_dri.so
usr/lib64/dri/rockchip_dri.so
usr/lib64/dri/rzg2l-du_dri.so
usr/lib64/dri/ssd130x_dri.so
usr/lib64/dri/st7586_dri.so
usr/lib64/dri/st7735r_dri.so
usr/lib64/dri/sti_dri.so
usr/lib64/dri/stm_dri.so
usr/lib64/dri/sun4i-drm_dri.so
usr/lib64/dri/swrast_dri.so
usr/lib64/dri/tegra_dri.so
usr/lib64/dri/udl_dri.so
usr/lib64/dri/virtio_gpu_dri.so
usr/lib64/dri/zink_dri.so
usr/lib64/dri/zynqmp-dpsub_dri.so
usr/lib64/libEGL_mesa.so.0
usr/lib64/libEGL_mesa.so.0.0.0
usr/lib64/libGLU.so.1
usr/lib64/libGLU.so.1.3.1
usr/lib64/libGLX_mesa.so.0
usr/lib64/libGLX_mesa.so.0.0.0
usr/lib64/libGLX_system.so.0
usr/lib64/libgbm.so.1
usr/lib64/libgbm.so.1.0.0
usr/lib64/libglapi.so.0
usr/lib64/libglapi.so.0.0.0
usr/lib64/libglut.so.3
usr/lib64/libglut.so.3.11.0
%%PORTDOCS%%usr/share/doc/freeglut/AUTHORS
%%PORTDOCS%%usr/share/doc/freeglut/COPYING
%%PORTDOCS%%usr/share/doc/freeglut/ChangeLog
%%PORTDOCS%%usr/share/doc/freeglut/README
%%PORTDOCS%%usr/share/doc/freeglut/download.html
%%PORTDOCS%%usr/share/doc/freeglut/freeglut.html
%%PORTDOCS%%usr/share/doc/freeglut/freeglut_logo.png
%%PORTDOCS%%usr/share/doc/freeglut/freeglut_user_interface.html
%%PORTDOCS%%usr/share/doc/freeglut/index.html
%%PORTDOCS%%usr/share/doc/freeglut/ogl_sm.png
%%PORTDOCS%%usr/share/doc/freeglut/progress.html
%%PORTDOCS%%usr/share/doc/freeglut/structure.html
%%DATADIR%%rc.d/00-mesa-defaults.conf
usr/share/glvnd/egl_vendor.d/50_mesa.json
