# sysutils/bacula/pkg-plist
etc/bacula-dir.conf
etc/bacula-fd.conf
etc/bacula-sd.conf
etc/console.conf
etc/query.sql
etc/rc.d/bacula.sh.sample
etc/rc.d/fd.sh.sample
sbin/bacula-dir
sbin/bacula-fd
sbin/bacula-sd
sbin/bextract
sbin/bls
sbin/bscan
sbin/btape
sbin/console
sbin/dbcheck
sbin/smtp
share/bacula/create_%%DBTYPE%%_database
share/bacula/delete_catalog_backup
share/bacula/drop_bacula_tables
share/bacula/drop_%%DBTYPE%%_tables
share/bacula/make_bacula_tables
share/bacula/make_%%DBTYPE%%_tables
share/bacula/make_catalog_backup
%%PORTDOCS%%share/doc/bacula/bacula.pdf
%%PORTDOCS%%share/doc/bacula/html-manual/1
%%PORTDOCS%%share/doc/bacula/html-manual/.cvsignore
%%PORTDOCS%%share/doc/bacula/html-manual/Bacula.apf
%%PORTDOCS%%share/doc/bacula/html-manual/add-nav.list
%%PORTDOCS%%share/doc/bacula/html-manual/add-nav.pl
%%PORTDOCS%%share/doc/bacula/html-manual/autochangers.html
%%PORTDOCS%%share/doc/bacula/html-manual/autochangers.wml
%%PORTDOCS%%share/doc/bacula/html-manual/back.gif
%%PORTDOCS%%share/doc/bacula/html-manual/bacula-features.html
%%PORTDOCS%%share/doc/bacula/html-manual/bacula-features.wml
%%PORTDOCS%%share/doc/bacula/html-manual/bacula-logo.gif
%%PORTDOCS%%share/doc/bacula/html-manual/bacula-styles.css
%%PORTDOCS%%share/doc/bacula/html-manual/bacula.book
%%PORTDOCS%%share/doc/bacula/html-manual/bacula.html
%%PORTDOCS%%share/doc/bacula/html-manual/bacula.wml
%%PORTDOCS%%share/doc/bacula/html-manual/bacula_template.inc
%%PORTDOCS%%share/doc/bacula/html-manual/bclogo.gif
%%PORTDOCS%%share/doc/bacula/html-manual/bootstrap.html
%%PORTDOCS%%share/doc/bacula/html-manual/bootstrap.wml
%%PORTDOCS%%share/doc/bacula/html-manual/bugs.html
%%PORTDOCS%%share/doc/bacula/html-manual/bugs.wml
%%PORTDOCS%%share/doc/bacula/html-manual/catalog.html
%%PORTDOCS%%share/doc/bacula/html-manual/catalog.wml
%%PORTDOCS%%share/doc/bacula/html-manual/catmaintenance.html
%%PORTDOCS%%share/doc/bacula/html-manual/catmaintenance.wml
%%PORTDOCS%%share/doc/bacula/html-manual/configure.html
%%PORTDOCS%%share/doc/bacula/html-manual/configure.wml
%%PORTDOCS%%share/doc/bacula/html-manual/console.html
%%PORTDOCS%%share/doc/bacula/html-manual/console.wml
%%PORTDOCS%%share/doc/bacula/html-manual/consoleconf.html
%%PORTDOCS%%share/doc/bacula/html-manual/consoleconf.wml
%%PORTDOCS%%share/doc/bacula/html-manual/daemonprotocol.html
%%PORTDOCS%%share/doc/bacula/html-manual/daemonprotocol.wml
%%PORTDOCS%%share/doc/bacula/html-manual/developers.html
%%PORTDOCS%%share/doc/bacula/html-manual/developers.wml
%%PORTDOCS%%share/doc/bacula/html-manual/dirdconf.html
%%PORTDOCS%%share/doc/bacula/html-manual/dirdconf.wml
%%PORTDOCS%%share/doc/bacula/html-manual/director.html
%%PORTDOCS%%share/doc/bacula/html-manual/director.wml
%%PORTDOCS%%share/doc/bacula/html-manual/disk.html
%%PORTDOCS%%share/doc/bacula/html-manual/disk.wml
%%PORTDOCS%%share/doc/bacula/html-manual/error.gif
%%PORTDOCS%%share/doc/bacula/html-manual/faq.html
%%PORTDOCS%%share/doc/bacula/html-manual/faq.wml
%%PORTDOCS%%share/doc/bacula/html-manual/file.gif
%%PORTDOCS%%share/doc/bacula/html-manual/file.html
%%PORTDOCS%%share/doc/bacula/html-manual/file.wml
%%PORTDOCS%%share/doc/bacula/html-manual/filedconf.html
%%PORTDOCS%%share/doc/bacula/html-manual/filedconf.wml
%%PORTDOCS%%share/doc/bacula/html-manual/flow.gif
%%PORTDOCS%%share/doc/bacula/html-manual/flow.jpeg
%%PORTDOCS%%share/doc/bacula/html-manual/gpl.html
%%PORTDOCS%%share/doc/bacula/html-manual/gpl.wml
%%PORTDOCS%%share/doc/bacula/html-manual/home.gif
%%PORTDOCS%%share/doc/bacula/html-manual/idle.gif
%%PORTDOCS%%share/doc/bacula/html-manual/imgdot-1x1-000000.gif
%%PORTDOCS%%share/doc/bacula/html-manual/imgdot-1x1-transp-ffffff.gif
%%PORTDOCS%%share/doc/bacula/html-manual/index.html
%%PORTDOCS%%share/doc/bacula/html-manual/index.wml
%%PORTDOCS%%share/doc/bacula/html-manual/install.html
%%PORTDOCS%%share/doc/bacula/html-manual/install.wml
%%PORTDOCS%%share/doc/bacula/html-manual/internaldb.html
%%PORTDOCS%%share/doc/bacula/html-manual/internaldb.wml
%%PORTDOCS%%share/doc/bacula/html-manual/internals.html
%%PORTDOCS%%share/doc/bacula/html-manual/internals.wml
%%PORTDOCS%%share/doc/bacula/html-manual/kaboom.html
%%PORTDOCS%%share/doc/bacula/html-manual/kaboom.wml
%%PORTDOCS%%share/doc/bacula/html-manual/lesser.html
%%PORTDOCS%%share/doc/bacula/html-manual/lesser.wml
%%PORTDOCS%%share/doc/bacula/html-manual/license.html
%%PORTDOCS%%share/doc/bacula/html-manual/license.wml
%%PORTDOCS%%share/doc/bacula/html-manual/md5.gif
%%PORTDOCS%%share/doc/bacula/html-manual/md5.html
%%PORTDOCS%%share/doc/bacula/html-manual/md5.wml
%%PORTDOCS%%share/doc/bacula/html-manual/mediaformat.html
%%PORTDOCS%%share/doc/bacula/html-manual/mediaformat.wml
%%PORTDOCS%%share/doc/bacula/html-manual/mempool.html
%%PORTDOCS%%share/doc/bacula/html-manual/mempool.wml
%%PORTDOCS%%share/doc/bacula/html-manual/messagesres.html
%%PORTDOCS%%share/doc/bacula/html-manual/messagesres.wml
%%PORTDOCS%%share/doc/bacula/html-manual/mysql.html
%%PORTDOCS%%share/doc/bacula/html-manual/mysql.wml
%%PORTDOCS%%share/doc/bacula/html-manual/netprotocol.html
%%PORTDOCS%%share/doc/bacula/html-manual/netprotocol.wml
%%PORTDOCS%%share/doc/bacula/html-manual/next.gif
%%PORTDOCS%%share/doc/bacula/html-manual/philosophical-gnu-sm.jpg
%%PORTDOCS%%share/doc/bacula/html-manual/porting.html
%%PORTDOCS%%share/doc/bacula/html-manual/porting.wml
%%PORTDOCS%%share/doc/bacula/html-manual/progs.html
%%PORTDOCS%%share/doc/bacula/html-manual/progs.wml
%%PORTDOCS%%share/doc/bacula/html-manual/projects.html
%%PORTDOCS%%share/doc/bacula/html-manual/projects.wml
%%PORTDOCS%%share/doc/bacula/html-manual/quickstart.html
%%PORTDOCS%%share/doc/bacula/html-manual/quickstart.wml
%%PORTDOCS%%share/doc/bacula/html-manual/recycling.html
%%PORTDOCS%%share/doc/bacula/html-manual/recycling.wml
%%PORTDOCS%%share/doc/bacula/html-manual/rescue.html
%%PORTDOCS%%share/doc/bacula/html-manual/rescue.wml
%%PORTDOCS%%share/doc/bacula/html-manual/restore.html
%%PORTDOCS%%share/doc/bacula/html-manual/restore.wml
%%PORTDOCS%%share/doc/bacula/html-manual/running.gif
%%PORTDOCS%%share/doc/bacula/html-manual/running.html
%%PORTDOCS%%share/doc/bacula/html-manual/running.wml
%%PORTDOCS%%share/doc/bacula/html-manual/security.html
%%PORTDOCS%%share/doc/bacula/html-manual/security.wml
%%PORTDOCS%%share/doc/bacula/html-manual/smartall.gif
%%PORTDOCS%%share/doc/bacula/html-manual/smartall.html
%%PORTDOCS%%share/doc/bacula/html-manual/smartall.wml
%%PORTDOCS%%share/doc/bacula/html-manual/sqlite.html
%%PORTDOCS%%share/doc/bacula/html-manual/sqlite.wml
%%PORTDOCS%%share/doc/bacula/html-manual/state.html
%%PORTDOCS%%share/doc/bacula/html-manual/state.wml
%%PORTDOCS%%share/doc/bacula/html-manual/status.gif
%%PORTDOCS%%share/doc/bacula/html-manual/storage.html
%%PORTDOCS%%share/doc/bacula/html-manual/storage.wml
%%PORTDOCS%%share/doc/bacula/html-manual/storedconf.html
%%PORTDOCS%%share/doc/bacula/html-manual/storedconf.wml
%%PORTDOCS%%share/doc/bacula/html-manual/tapetesting.html
%%PORTDOCS%%share/doc/bacula/html-manual/tapetesting.wml
%%PORTDOCS%%share/doc/bacula/html-manual/thanks.gif
%%PORTDOCS%%share/doc/bacula/html-manual/thanks.html
%%PORTDOCS%%share/doc/bacula/html-manual/thanks.wml
%%PORTDOCS%%share/doc/bacula/html-manual/tips.html
%%PORTDOCS%%share/doc/bacula/html-manual/tips.wml
%%PORTDOCS%%share/doc/bacula/html-manual/verify.html
%%PORTDOCS%%share/doc/bacula/html-manual/verify.wml
%%PORTDOCS%%share/doc/bacula/html-manual/waiting.gif
%%PORTDOCS%%share/doc/bacula/html-manual/win32.html
%%PORTDOCS%%share/doc/bacula/html-manual/win32.wml
%%PORTDOCS%%share/doc/bacula/html-manual/wininstall-1.gif
%%PORTDOCS%%share/doc/bacula/html-manual/wininstall-2.gif
%%PORTDOCS%%share/doc/bacula/html-manual/wininstall-3.gif
%%PORTDOCS%%share/doc/bacula/html-manual/wininstall-4.gif
%%PORTDOCS%%share/doc/bacula/html-manual/wininstall-5.gif
%%PORTDOCS%%share/doc/bacula/html-manual/wininstall-6.gif
%%PORTDOCS%%share/doc/bacula/html-manual/wininstall-7.gif
%%PORTDOCS%%share/doc/bacula/html-manual/wininstall-8.gif
%%PORTDOCS%%share/doc/bacula/html-manual/wininstall-9.gif
%%PORTDOCS%%share/doc/bacula/html-manual/working.gif
%%PORTDOCS%%@dirrm share/doc/bacula/html-manual
%%PORTDOCS%%@dirrm share/doc/bacula
@dirrm share/bacula
