# net/freeswitch-core/files/freeswitch.in
#!/bin/sh

# PROVIDE: freeswitch
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable freeswitch:
#
# freeswitch_enable="YES"
#
#
#  freeswitch_user="freeswitch"
#  freeswitch_group="freeswitch"

. /etc/rc.subr

name=freeswitch
rcvar=`set_rcvar`

load_rc_config $name

freeswitch_enable=${freeswitch_enable:-"NO"}
freeswitch_user=${freeswitch_user:-"freeswitch"}
freeswitch_group=${freeswitch_group:-"freeswitch"}
freeswitch_flags=${freeswitch_flags:-"-nc -waste"}

command=%%PREFIX%%/bin/freeswitch
command_args="-u ${freeswitch_user} -g ${freeswitch_group}"

pidfile=${freeswitch_pidfile:-"/var/run/freeswitch/freeswitch.pid"}

stop_precmd="freeswitch_stop"

freeswitch_stop () {
        echo "Stopping FreeSWITCH."
        %%PREFIX%%/bin/freeswitch -stop
	wait_for_pids `cat $pidfile`
}

run_rc_command "$1"
