# net/xrdp/pkg-plist
bin/xrdp-dis
bin/xrdp-dumpfv1
bin/xrdp-genkeymap
bin/xrdp-keygen
bin/xrdp-sesadmin
bin/xrdp-sesrun
etc/pam.d/xrdp-sesman
etc/rc.d/xrdp
etc/rc.d/xrdp-sesman
%%ETCDIR%%/km-00000405.ini
%%ETCDIR%%/km-00000406.ini
%%ETCDIR%%/km-00000407.ini
%%ETCDIR%%/km-00000409.ini
%%ETCDIR%%/km-0000040a.ini
%%ETCDIR%%/km-0000040b.ini
%%ETCDIR%%/km-0000040c.ini
%%ETCDIR%%/km-0000040e.ini
%%ETCDIR%%/km-00000410.ini
%%ETCDIR%%/km-00000411.ini
%%ETCDIR%%/km-00000412.ini
%%ETCDIR%%/km-00000414.ini
%%ETCDIR%%/km-00000415.ini
%%ETCDIR%%/km-00000416.ini
%%ETCDIR%%/km-00000419.ini
%%ETCDIR%%/km-0000041d.ini
%%ETCDIR%%/km-00000426.ini
%%ETCDIR%%/km-00010426.ini
%%ETCDIR%%/km-00020426.ini
%%ETCDIR%%/km-00000807.ini
%%ETCDIR%%/km-00000809.ini
%%ETCDIR%%/km-0000080a.ini
%%ETCDIR%%/km-0000080c.ini
%%ETCDIR%%/km-00000813.ini
%%ETCDIR%%/km-00000816.ini
%%ETCDIR%%/km-0000100c.ini
%%ETCDIR%%/km-00010409.ini
%%ETCDIR%%/km-19360409.ini
%%ETCDIR%%/pulse/default.pa
@sample %%ETCDIR%%/gfx.toml.sample
@sample %%ETCDIR%%/reconnectwm.sh.sample
@sample %%ETCDIR%%/sesman.ini.sample
@sample %%ETCDIR%%/startwm.sh.sample
@sample %%ETCDIR%%/xrdp_keyboard.ini.sample
@sample %%ETCDIR%%/xrdp.ini.sample
include/ms-erref.h
include/ms-fscc.h
include/ms-rdpbcgr.h
include/ms-rdpeclip.h
include/ms-rdpedisp.h
include/ms-rdpefs.h
include/ms-rdpegdi.h
include/ms-rdpele.h
include/ms-rdperp.h
include/ms-smb2.h
include/painter.h
include/rfxcodec_common.h
include/rfxcodec_decode.h
include/rfxcodec_encode.h
include/xrdp_client_info.h
include/xrdp_constants.h
include/xrdp_rail.h
include/xrdp_sockets.h
lib/libpainter.a
lib/librfxencode.a
lib/xrdp/libcommon.a
lib/xrdp/libcommon.so
lib/xrdp/libcommon.so.0
lib/xrdp/libcommon.so.0.0.0
lib/xrdp/libipm.a
lib/xrdp/libipm.so
lib/xrdp/libipm.so.0
lib/xrdp/libipm.so.0.0.0
lib/xrdp/libmc.a
lib/xrdp/libmc.so
lib/xrdp/libsesman.a
lib/xrdp/libsesman.so
lib/xrdp/libsesman.so.0
lib/xrdp/libsesman.so.0.0.0
lib/xrdp/libtoml.a
lib/xrdp/libtoml.so
lib/xrdp/libtoml.so.1
lib/xrdp/libtoml.so.1.0.0
lib/xrdp/libvnc.a
lib/xrdp/libvnc.so
lib/xrdp/libxrdp.a
lib/xrdp/libxrdp.so
lib/xrdp/libxrdp.so.0
lib/xrdp/libxrdp.so.0.0.0
lib/xrdp/libxrdpapi.a
lib/xrdp/libxrdpapi.so
lib/xrdp/libxrdpapi.so.0
lib/xrdp/libxrdpapi.so.0.0.0
lib/xrdp/libxup.a
lib/xrdp/libxup.so
libdata/pkgconfig/libpainter.pc
libdata/pkgconfig/rfxcodec.pc
libdata/pkgconfig/xrdp.pc
libexec/xrdp/waitforx
libexec/xrdp/xrdp-droppriv
libexec/xrdp/xrdp-sesexec
sbin/xrdp
sbin/xrdp-chansrv
sbin/xrdp-sesman
share/man/man1/xrdp-dis.1.gz
share/man/man5/gfx.toml.5.gz
share/man/man5/sesman.ini.5.gz
share/man/man5/xrdp.ini.5.gz
share/man/man8/xrdp-chansrv.8.gz
share/man/man8/xrdp-dumpfv1.8.gz
share/man/man8/xrdp-genkeymap.8.gz
share/man/man8/xrdp-keygen.8.gz
share/man/man8/xrdp-sesadmin.8.gz
share/man/man8/xrdp-sesman.8.gz
share/man/man8/xrdp-sesrun.8.gz
share/man/man8/xrdp.8.gz
%%DATADIR%%/ad24b.bmp
%%DATADIR%%/ad256.bmp
%%DATADIR%%/cursor0.cur
%%DATADIR%%/cursor1.cur
%%DATADIR%%/README.logo
%%DATADIR%%/sans-10.fv1
%%DATADIR%%/sans-18.fv1
%%DATADIR%%/xrdp_logo.bmp
%%DATADIR%%/xrdp_logo.png
%%DATADIR%%/xrdp-chkpriv
%%DATADIR%%/xrdp24b.bmp
%%DATADIR%%/xrdp256.bmp
