# net/xrdp/pkg-plist
bin/xrdp-dis
bin/xrdp-genkeymap
bin/xrdp-keygen
bin/xrdp-sesadmin
bin/xrdp-sesrun
lib/libpainter.a
lib/libpainter.so
lib/libpainter.so.0
lib/libpainter.so.0.0.0
lib/librfxencode.so
lib/librfxencode.so.0.0.0
lib/librfxencode.a
lib/librfxencode.so.0
include/ms-erref.h
include/ms-fscc.h
include/ms-rdpbcgr.h
include/ms-rdpefs.h
include/ms-rdpegdi.h
include/ms-rdpele.h
include/ms-rdperp.h
include/ms-smb2.h
include/painter.h
include/rfxcodec_common.h
include/rfxcodec_decode.h
include/rfxcodec_encode.h
include/xrdp_client_info.h
include/xrdp_constants.h
include/xrdp_rail.h
include/xrdp_sockets.h
libdata/pkgconfig/libpainter.pc
libdata/pkgconfig/rfxcodec.pc
libdata/pkgconfig/xrdp.pc
lib/xrdp/libcommon.a
lib/xrdp/libcommon.so
lib/xrdp/libcommon.so.0
lib/xrdp/libcommon.so.0.0.0
lib/xrdp/libmc.a
lib/xrdp/libmc.so
lib/xrdp/libscp.a
lib/xrdp/libscp.so
lib/xrdp/libscp.so.0
lib/xrdp/libscp.so.0.0.0
lib/xrdp/libvnc.a
lib/xrdp/libvnc.so
lib/xrdp/libxrdp.a
lib/xrdp/libxrdp.so
lib/xrdp/libxrdp.so.0
lib/xrdp/libxrdp.so.0.0.0
lib/xrdp/libxup.a
lib/xrdp/libxup.so
lib/xrdp/libxrdpapi.so
lib/xrdp/libxrdpapi.a
lib/xrdp/libxrdpapi.so.0
lib/xrdp/libxrdpapi.so.0.0.0
etc/pam.d/xrdp-sesman
etc/rc.d/xrdp
etc/rc.d/xrdp-sesman
%%ETCDIR%%/km-00000406.ini
%%ETCDIR%%/km-00000407.ini
%%ETCDIR%%/km-00000409.ini
%%ETCDIR%%/km-00010409.ini
%%ETCDIR%%/km-0000040a.ini
%%ETCDIR%%/km-0000040b.ini
%%ETCDIR%%/km-0000040c.ini
%%ETCDIR%%/km-00000410.ini
%%ETCDIR%%/km-00000411.ini
%%ETCDIR%%/km-00000412.ini
%%ETCDIR%%/km-00000414.ini
%%ETCDIR%%/km-00000415.ini
%%ETCDIR%%/km-00000416.ini
%%ETCDIR%%/km-00000419.ini
%%ETCDIR%%/km-0000041d.ini
%%ETCDIR%%/km-00000807.ini
%%ETCDIR%%/km-00000809.ini
%%ETCDIR%%/km-0000080a.ini
%%ETCDIR%%/km-0000080c.ini
%%ETCDIR%%/km-00000813.ini
%%ETCDIR%%/km-00000816.ini
%%ETCDIR%%/km-0000100c.ini
%%ETCDIR%%/pulse/default.pa
@postunexec if cmp -s %D/etc/xrdp/sesman.ini.sample %D/etc/xrdp/sesman.ini; then rm -f %D/etc/xrdp/sesman.ini; fi
@postunexec if cmp -s %D/etc/xrdp/startwm.sh.sample %D/etc/xrdp/startwm.sh; then rm -f %D/etc/xrdp/startwm.sh; fi
@postunexec if cmp -s %D/etc/xrdp/reconnectwm.sh.sample %D/etc/xrdp/reconnectwm.sh; then rm -f %D/etc/xrdp/reconnectwm.sh; fi
@postunexec if cmp -s %D/etc/xrdp/xrdp.ini.sample %D/etc/xrdp/xrdp.ini; then rm -f %D/etc/xrdp/xrdp.ini; fi
@postunexec if cmp -s %D/etc/xrdp/xrdp_keyboard.ini.sample %D/etc/xrdp/xrdp_keyboard.ini; then rm -f %D/etc/xrdp/xrdp_keyboard.ini; fi
@sample %%ETCDIR%%/sesman.ini.sample
@sample %%ETCDIR%%/startwm.sh.sample
@sample %%ETCDIR%%/reconnectwm.sh.sample
@sample %%ETCDIR%%/xrdp.ini.sample
@sample %%ETCDIR%%/xrdp_keyboard.ini.sample
man/man1/xrdp-dis.1.gz
man/man5/sesman.ini.5.gz
man/man8/xrdp.8.gz
man/man5/xrdp.ini.5.gz
man/man8/xrdp-chansrv.8.gz
man/man8/xrdp-genkeymap.8.gz
man/man8/xrdp-keygen.8.gz
man/man8/xrdp-sesadmin.8.gz
man/man8/xrdp-sesman.8.gz
man/man8/xrdp-sesrun.8.gz
sbin/xrdp
sbin/xrdp-chansrv
sbin/xrdp-sesman
%%DATADIR%%/ad24b.bmp
%%DATADIR%%/ad256.bmp
%%DATADIR%%/cursor0.cur
%%DATADIR%%/cursor1.cur
%%DATADIR%%/sans-10.fv1
%%DATADIR%%/xrdp24b.bmp
%%DATADIR%%/xrdp256.bmp
%%DATADIR%%/xrdp_logo.bmp
