# security/heimdal/pkg-plist
bin/afslog
bin/ftp
bin/kauth
bin/kdestroy
bin/kf
bin/kgetcred
bin/kinit
bin/klist
bin/kpasswd
bin/krb5-config
bin/login
bin/mk_cmds
bin/otp
bin/otpprint
bin/pagsh
bin/pfrom
bin/rcp
bin/rsh
bin/string2key
bin/su
bin/telnet
bin/verify_krb5_conf
etc/rc.d/kdc.sh.sample
include/asn1_err.h
include/base64.h
include/der.h
include/editline.h
include/getarg.h
include/gssapi.h
include/hdb-private.h
include/hdb-protos.h
include/hdb.h
include/hdb_asn1.h
include/hdb_err.h
include/heim_err.h
include/k524_err.h
include/kadm5/admin.h
include/kadm5/kadm5-private.h
include/kadm5/kadm5-protos.h
include/kadm5/kadm5_err.h
include/kadm5/private.h
include/kafs.h
include/krb5-private.h
include/krb5-protos.h
include/krb5-types.h
include/krb5.h
include/krb5_asn1.h
include/krb5_err.h
include/otp.h
include/parse_bytes.h
include/parse_time.h
include/parse_units.h
include/resolve.h
include/roken-common.h
include/roken.h
include/rtbl.h
include/sl.h
include/ss/ss.h
include/xdbm.h
info/dir
lib/libasn1.a
lib/libasn1.so
lib/libasn1.so.6
lib/libeditline.a
lib/libgssapi.a
lib/libgssapi.so
lib/libgssapi.so.5
lib/libhdb.a
lib/libhdb.so
lib/libhdb.so.7
lib/libkadm5clnt.a
lib/libkadm5clnt.so
lib/libkadm5clnt.so.6
lib/libkadm5srv.a
lib/libkadm5srv.so
lib/libkadm5srv.so.7
lib/libkafs.a
lib/libkafs.so
lib/libkafs.so.4
lib/libkrb5.a
lib/libkrb5.so
lib/libkrb5.so.20
lib/libotp.a
lib/libotp.so
lib/libotp.so.1
lib/libroken.a
lib/libroken.so
lib/libroken.so.16
lib/libsl.a
lib/libsl.so
lib/libsl.so.1
lib/libss.a
lib/libss.so
lib/libss.so.1
libexec/ftpd
libexec/hprop
libexec/hpropd
libexec/ipropd-master
libexec/ipropd-slave
libexec/kadmind
libexec/kdc
libexec/kfd
libexec/kpasswdd
libexec/popper
libexec/push
libexec/rshd
libexec/telnetd
sbin/dump_log
sbin/kadmin
sbin/kstash
sbin/ktutil
sbin/replay_log
sbin/truncate_log
@dirrm include/kadm5
@dirrm include/ss
