# security/heimdal/pkg-plist
bin/afslog
bin/ftp
bin/gss
bin/hxtool
bin/idn-lookup
bin/kauth
bin/kdestroy
bin/kf
bin/kgetcred
bin/kinit
bin/klist
bin/kpasswd
bin/krb5-config
bin/kswitch
bin/login
bin/otp
bin/otpprint
bin/pagsh
bin/pfrom
bin/rcp
bin/rsh
bin/string2key
bin/su
bin/telnet
bin/verify_krb5_conf
etc/rc.d/kdc.sh.sample
include/asn1-common.h
include/asn1_err.h
include/base64.h
include/cms_asn1.h
include/crmf_asn1.h
include/der-private.h
include/der-protos.h
include/der.h
include/digest_asn1.h
include/getarg.h
include/gssapi.h
include/gssapi/gkrb5_err.h
include/gssapi/gssapi.h
include/gssapi/gssapi_ntlm.h
include/gssapi/gssapi_krb5.h
include/gssapi/gssapi_spnego.h
%%HCRYPTO%%include/hcrypto/aes.h
%%HCRYPTO%%include/hcrypto/bn.h
%%HCRYPTO%%include/hcrypto/des.h
%%HCRYPTO%%include/hcrypto/dh.h
%%HCRYPTO%%include/hcrypto/dsa.h
%%HCRYPTO%%include/hcrypto/ec.h
%%HCRYPTO%%include/hcrypto/ecdh.h
%%HCRYPTO%%include/hcrypto/ecdsa.h
%%HCRYPTO%%include/hcrypto/engine.h
%%HCRYPTO%%include/hcrypto/evp.h
%%HCRYPTO%%include/hcrypto/evp-cc.h
%%HCRYPTO%%include/hcrypto/evp-hcrypto.h
%%HCRYPTO%%include/hcrypto/hmac.h
%%HCRYPTO%%include/hcrypto/md2.h
%%HCRYPTO%%include/hcrypto/md4.h
%%HCRYPTO%%include/hcrypto/md5.h
%%HCRYPTO%%include/hcrypto/pkcs12.h
%%HCRYPTO%%include/hcrypto/rand.h
%%HCRYPTO%%include/hcrypto/rc2.h
%%HCRYPTO%%include/hcrypto/rc4.h
%%HCRYPTO%%include/hcrypto/rsa.h
%%HCRYPTO%%include/hcrypto/sha.h
%%HCRYPTO%%include/hcrypto/ui.h
include/hdb-protos.h
include/hdb.h
include/hdb_asn1.h
include/hdb_err.h
include/heim_asn1.h
include/heim-ipc.h
include/heim_err.h
include/heimntlm-protos.h
include/heimntlm.h
include/hex.h
include/hx509-protos.h
include/hx509.h
include/hx509_err.h
include/k524_err.h
include/kadm5/admin.h
include/kadm5/kadm5-private.h
include/kadm5/kadm5-protos.h
include/kadm5/kadm5-pwcheck.h
include/kadm5/kadm5_err.h
include/kadm5/private.h
include/kafs.h
include/kdc-protos.h
include/kdc.h
include/krb5-private.h
include/krb5-protos.h
include/krb5-types.h
include/krb5.h
include/krb5/locate_plugin.h
include/krb5/send_to_kdc_plugin.h
include/krb5/windc_plugin.h
include/krb5_asn1.h
include/krb5_ccapi.h
include/krb5_err.h
include/kx509_asn1.h
include/ocsp_asn1.h
include/otp.h
include/parse_bytes.h
include/parse_time.h
include/parse_units.h
include/pkcs10_asn1.h
include/pkcs12_asn1.h
include/pkcs8_asn1.h
include/pkcs9_asn1.h
include/pkinit_asn1.h
include/resolve.h
include/rfc2459_asn1.h
include/roken-common.h
include/roken.h
include/rtbl.h
include/sl.h
include/wind.h
include/wind_err.h
include/xdbm.h
lib/libasn1.a
lib/libasn1.la
lib/libasn1.so
lib/libasn1.so.8
lib/libgssapi.a
lib/libgssapi.la
lib/libgssapi.so
lib/libgssapi.so.2
%%HCRYPTO%%lib/libhcrypto.a
%%HCRYPTO%%lib/libhcrypto.la
%%HCRYPTO%%lib/libhcrypto.so
%%HCRYPTO%%lib/libhcrypto.so.5
%%SQLITE%%lib/libheimsqlite.a
%%SQLITE%%lib/libheimsqlite.la
%%SQLITE%%lib/libheimsqlite.so
%%SQLITE%%lib/libheimsqlite.so.0
lib/libhdb.a
lib/libhdb.la
lib/libhdb.so
lib/libhdb.so.11
lib/libheimntlm.a
lib/libheimntlm.la
lib/libheimntlm.so
lib/libheimntlm.so.1
lib/libhx509.a
lib/libhx509.la
lib/libhx509.so
lib/libhx509.so.5
lib/libkadm5clnt.a
lib/libkadm5clnt.la
lib/libkadm5clnt.so
lib/libkadm5clnt.so.7
lib/libkadm5srv.a
lib/libkadm5srv.la
lib/libkadm5srv.so
lib/libkadm5srv.so.8
lib/libkafs.a
lib/libkafs.la
lib/libkafs.so
lib/libkafs.so.5
lib/libkdc.a
lib/libkdc.la
lib/libkdc.so
lib/libkdc.so.2
lib/libkrb5.a
lib/libkrb5.la
lib/libkrb5.so
lib/libkrb5.so.26
lib/libotp.a
lib/libotp.la
lib/libotp.so
lib/libotp.so.1
lib/libroken.a
lib/libroken.la
lib/libroken.so
lib/libroken.so.19
lib/libsl.a
lib/libsl.la
lib/libsl.so
lib/libsl.so.2
lib/libwind.a
lib/libwind.la
lib/libwind.so
lib/libwind.so.0
lib/windc.a
lib/windc.la
lib/windc.so
lib/windc.so.0
libdata/pkgconfig/heimdal-gssapi.pc
libexec/digest-service
libexec/ftpd
libexec/heimdal/asn1_compile
libexec/heimdal/asn1_print
libexec/heimdal/slc
libexec/hprop
libexec/hpropd
libexec/ipropd-master
libexec/ipropd-slave
libexec/kadmind
%%KCM%%libexec/kcm
libexec/kdc
libexec/kdigest
libexec/kfd
libexec/kimpersonate
libexec/kpasswdd
libexec/popper
libexec/push
libexec/rshd
libexec/telnetd
sbin/iprop-log
sbin/kadmin
sbin/kstash
sbin/ktutil
@dirrm include/krb5
@dirrm include/kadm5
%%HCRYPTO%%@dirrm include/hcrypto
@dirrm include/gssapi
@dirrm include/roken
@dirrm libexec/heimdal
