# security/prelude-lml/files/prelude-lml.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: prelude_lml
# REQUIRE: LOGIN postgresql mysql
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# prelude_lml_enable (bool):   Set to NO by default.
#                              Set it to YES to enable Prelude LML.
# prelude_lml_flags  (flags):  Set arguments to use with Prelude LML

. /etc/rc.subr

name="prelude_lml"
rcvar=prelude_lml_enable

load_rc_config $name

: ${prelude_lml_enable:="NO"}

command=%%PREFIX%%/bin/prelude-lml
pidfile=/var/run/prelude-lml.pid
command_args="-d -P $pidfile"

run_rc_command "$1"
