# sysutils/bacula11-server/files/bacula-fd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: bacula_fd
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# bacula_fd_enable  (bool):  Set to NO by default.
#               Set it to YES to enable bacula_fd.
# bacula_fd_flags (params):  Set params used to start bacula_fd.
#

. /etc/rc.subr

name="bacula_fd"
rcvar=${name}_enable
command=%%PREFIX%%/sbin/bacula-fd

load_rc_config $name

: ${bacula_fd_enable="NO"}
: ${bacula_fd_flags=" -u root -g wheel -v -c %%PREFIX%%/etc/bacula/bacula-fd.conf"}
: ${bacula_fd_pidfile="/var/run/bacula-fd.9102.pid"}

pidfile="${bacula_fd_pidfile}"

run_rc_command "$1"
