# sysutils/rsyncmanager/pkg-plist
bin/rsyncmanager.rb
etc/rc.d/rsyncmanager.sh
etc/rsyncmanager.xml.sample
