# net/gnu-radius/pkg-plist
bin/builddbm
bin/radgrep
bin/radlast
bin/radping
bin/radtest
bin/radwho
bin/radzap
etc/raddb/access.deny
etc/raddb/client.conf
etc/raddb/clients
etc/raddb/config
etc/raddb/dict/ascend
etc/raddb/dict/ascend.vs
etc/raddb/dict/cisco
etc/raddb/dict/compat
etc/raddb/dict/livingston
etc/raddb/dict/standard
etc/raddb/dict/tunnel
etc/raddb/dict/usr
etc/raddb/dict/voip
etc/raddb/dictionary
etc/raddb/hints
etc/raddb/huntgroups
etc/raddb/naslist
etc/raddb/nastypes
etc/raddb/nas.rc
etc/raddb/realms
etc/raddb/sqlserver
etc/raddb/users
etc/rc.d/gnu-radius.sh
@unexec install-info --delete %D/info/radius.info %D/info/dir
info/radius.info
info/radius.info-1
info/radius.info-2
info/radius.info-3
info/radius.info-4
info/radius.info-5
info/radius.info-6
info/radius.info-7
@exec install-info %D/info/radius.info %D/info/dir
sbin/radauth
sbin/radctl
sbin/radiusd
share/locale/es/LC_MESSAGES/radius.mo
share/locale/pl/LC_MESSAGES/radius.mo
share/locale/ru/LC_MESSAGES/radius.mo
%%GUILE%%share/radius/guile-procedures.txt
share/radius/1.2/scheme/boot.scm
share/radius/1.2/scheme/ttl.scm
share/radius/1.2/rewrite/checknas.rw
share/radius/1.2/rewrite/log-hook.rw
@dirrm etc/raddb/dict
@dirrm etc/raddb
@unexec rm -f /var/run/radiusd.pid
@unexec rm -rf /var/log/radius
@unexec echo "Leaving /var/run/radutmp, /var/log/radwtmp and /var/log/radstat alone."
