# net/gnu-radius/pkg-plist
@unexec %D/etc/rc.d/radiusd.sh forcestop 2>/dev/null || true
bin/builddbm
bin/radgrep
bin/radlast
bin/radping
bin/radtest
bin/radwho
bin/radzap
sbin/radauth
sbin/radctl
sbin/radiusd
@unexec if cmp -s %D/etc/raddb/access.deny %D/etc/raddb/access.deny-dist; then rm -f %D/etc/raddb/access.deny; fi
etc/raddb/access.deny-dist
@unexec if cmp -s %D/etc/raddb/client.conf %D/etc/raddb/client.conf-dist; then rm -f %D/etc/raddb/client.conf; fi
etc/raddb/client.conf-dist
@unexec if cmp -s %D/etc/raddb/clients %D/etc/raddb/clients-dist; then rm -f %D/etc/raddb/clients; fi
etc/raddb/clients-dist
@unexec if cmp -s %D/etc/raddb/config %D/etc/raddb/config-dist; then rm -f %D/etc/raddb/config; fi
etc/raddb/config-dist
etc/raddb/dict/ascend
etc/raddb/dict/ascend.vs
etc/raddb/dict/cisco
etc/raddb/dict/compat
etc/raddb/dict/livingston
etc/raddb/dict/microsoft
etc/raddb/dict/standard
etc/raddb/dict/tunnel
etc/raddb/dict/usr
etc/raddb/dict/voip
@unexec if cmp -s %D/etc/raddb/dictionary %D/etc/raddb/dictionary-dist; then rm -f %D/etc/raddb/dictionary; fi
etc/raddb/dictionary-dist
@unexec if cmp -s %D/etc/raddb/hints %D/etc/raddb/hints-dist; then rm -f %D/etc/raddb/hints; fi
etc/raddb/hints-dist
@unexec if cmp -s %D/etc/raddb/huntgroups %D/etc/raddb/huntgroups-dist; then rm -f %D/etc/raddb/huntgroups; fi
etc/raddb/huntgroups-dist
etc/raddb/nas.rc
@unexec if cmp -s %D/etc/raddb/naslist %D/etc/raddb/naslist-dist; then rm -f %D/etc/raddb/naslist; fi
etc/raddb/naslist-dist
@unexec if cmp -s %D/etc/raddb/nastypes %D/etc/raddb/nastypes-dist; then rm -f %D/etc/raddb/nastypes; fi
etc/raddb/nastypes-dist
@unexec if cmp -s %D/etc/raddb/realms %D/etc/raddb/realms-dist; then rm -f %D/etc/raddb/realms; fi
etc/raddb/realms-dist
@unexec if cmp -s %D/etc/raddb/sqlserver %D/etc/raddb/sqlserver-dist; then rm -f %D/etc/raddb/sqlserver; fi
etc/raddb/sqlserver-dist
@unexec if cmp -s %D/etc/raddb/users %D/etc/raddb/users-dist; then rm -f %D/etc/raddb/users; fi
etc/raddb/users-dist
etc/rc.d/radiusd.sh
include/radius/argcv.h
include/radius/debug.h
include/radius/dictionary.h
include/radius/envar.h
include/radius/list.h
include/radius/md4.h
include/radius/md5.h
include/radius/mem.h
include/radius/raddbm.h
include/radius/radius.h
include/radius/radpaths.h
include/radius/radscm.h
include/radius/radutmp.h
include/radius/sha1.h
include/radius/symtab.h
include/radius/types.h
lib/libgnuradius.a
lib/libgnuradius.so
lib/libgnuradius.la
lib/libgnuradius.so.0
share/aclocal/radius.m4
%%CLIENT%%bin/radsession
%%CLIENT%%lib/libguile-gnuradius-v-%%PORTVERSION%%.so
%%CLIENT%%lib/libradscm.a
%%CLIENT%%lib/libradscm.la
%%CLIENT%%lib/libradscm.so
%%CLIENT%%lib/libradscm.so.1
%%CLIENT%%lib/libservscm.a
%%CLIENT%%lib/security/pam_radius.a
%%CLIENT%%lib/security/pam_radius.la
%%CLIENT%%lib/security/pam_radius.so
%%CLIENT%%libexec/nas.scm
%%NLS%%share/locale/ca/LC_MESSAGES/radius.mo
%%NLS%%share/locale/es/LC_MESSAGES/radius.mo
%%NLS%%share/locale/fr/LC_MESSAGES/radius.mo
%%NLS%%share/locale/id/LC_MESSAGES/radius.mo
%%NLS%%share/locale/pl/LC_MESSAGES/radius.mo
%%NLS%%share/locale/ru/LC_MESSAGES/radius.mo
@comment %%NLS%%share/locale/rw/LC_MESSAGES/radius.mo
%%NLS%%share/locale/uk/LC_MESSAGES/radius.mo
%%NLS%%share/locale/vi/LC_MESSAGES/radius.mo
%%CLIENT%%%%DATADIR%%/%%PORTVERSION%%/scheme/gnuradius.scm
%%CLIENT%%%%DATADIR%%/%%PORTVERSION%%/scheme/ipalloc.scm
%%CLIENT%%%%DATADIR%%/%%PORTVERSION%%/scheme/radiusd.scm
%%CLIENT%%%%DATADIR%%/%%PORTVERSION%%/scheme/ttl.scm
%%CLIENT%%%%DATADIR%%/%%PORTVERSION%%/scheme/guile-procedures.txt
%%DATADIR%%/%%PORTVERSION%%/rewrite/checknas.rw
%%DATADIR%%/%%PORTVERSION%%/rewrite/log-hook.rw
%%DATADIR%%/%%PORTVERSION%%/rewrite/nas-ip.rw
%%MYSQL%%lib/radius/%%PORTVERSION%%/modules/mysql.a
%%MYSQL%%lib/radius/%%PORTVERSION%%/modules/mysql.la
%%MYSQL%%lib/radius/%%PORTVERSION%%/modules/mysql.so
%%PGSQL%%lib/radius/%%PORTVERSION%%/modules/postgres.a
%%PGSQL%%lib/radius/%%PORTVERSION%%/modules/postgres.la
%%PGSQL%%lib/radius/%%PORTVERSION%%/modules/postgres.so
%%EMACS%%%%EMACS_SITE_LISPDIR%%/radconf-mode.el
%%EMACS%%%%EMACS_SITE_LISPDIR%%/radius-mode.el
%%EMACS%%%%EMACS_SITE_LISPDIR%%/rewrite-mode.el
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/README
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/ascend.rw
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/c3620.rw
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/cisco.rw
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/comos.rw
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/config.syntax
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/filter.rw
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/hints-0.96
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/jetstream.rw
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/nt.rw
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/pam.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/radius.php
@dirrm etc/raddb/dict
@dirrmtry etc/raddb
@dirrm include/radius
@dirrm %%DATADIR%%/%%PORTVERSION%%/scheme
@dirrm %%DATADIR%%/%%PORTVERSION%%/rewrite
@dirrm %%DATADIR%%/%%PORTVERSION%%
@dirrm %%DATADIR%%
%%PORTEXAMPLES%%@dirrm %%EXAMPLESDIR%%
@dirrm lib/radius/%%PORTVERSION%%/modules
@dirrm lib/radius/%%PORTVERSION%%
@dirrm lib/radius
@dirrm lib/security
@unexec rm -f %B/run/radiusd.pid
@unexec rm -rf %B/log/radius 2>&1 >/dev/null || true
@unexec echo "Leaving /var/run/radutmp, /var/log/radwtmp and /var/log/radstat alone."
