# net/gnu-radius/pkg-plist
@unexec %D/etc/rc.d/radiusd.sh forcestop 2>/dev/null || true
bin/builddbm
bin/radgrep
bin/radlast
bin/radping
bin/radtest
bin/radwho
bin/radzap
sbin/radauth
sbin/radctl
sbin/radiusd
etc/raddb/access.deny
etc/raddb/client.conf
etc/raddb/clients
etc/raddb/config
etc/raddb/dict/ascend
etc/raddb/dict/ascend.vs
etc/raddb/dict/cisco
etc/raddb/dict/compat
etc/raddb/dict/livingston
etc/raddb/dict/microsoft
etc/raddb/dict/standard
etc/raddb/dict/tunnel
etc/raddb/dict/usr
etc/raddb/dict/voip
etc/raddb/dictionary
etc/raddb/hints
etc/raddb/huntgroups
etc/raddb/nas.rc
etc/raddb/naslist
etc/raddb/nastypes
etc/raddb/realms
etc/raddb/sqlserver
etc/raddb/users
etc/rc.d/radiusd.sh
include/radius/argcv.h
include/radius/argp.h
include/radius/dictionary.h
include/radius/envar.h
include/radius/getopt.h
include/radius/list.h
include/radius/md4.h
include/radius/md5.h
include/radius/mem.h
include/radius/obstack.h
include/radius/radargp.h
include/radius/raddbm.h
include/radius/radius.h
include/radius/radpaths.h
include/radius/radscm.h
include/radius/radutmp.h
include/radius/sha1.h
include/radius/symtab.h
include/radius/types.h
lib/libgnuradius.a
lib/libgnuradius.so
lib/libgnuradius.so.0
%%GUILE%%bin/radsession
%%GUILE%%lib/libradscm.so
%%GUILE%%lib/libradscm.so.0
%%CLIENT%%lib/security/pam_radius.so
%%CLIENT%%libexec/nas.scm
share/locale/ca/LC_MESSAGES/radius.mo
share/locale/es/LC_MESSAGES/radius.mo
share/locale/pl/LC_MESSAGES/radius.mo
share/locale/ru/LC_MESSAGES/radius.mo
%%GUILE%%%%DATADIR%%/%%PORTVERSION%%/scheme/gnuradius.scm
%%GUILE%%%%DATADIR%%/%%PORTVERSION%%/scheme/ipalloc.scm
%%GUILE%%%%DATADIR%%/%%PORTVERSION%%/scheme/radiusd.scm
%%GUILE%%%%DATADIR%%/%%PORTVERSION%%/scheme/ttl.scm
%%GUILE%%%%DATADIR%%/%%PORTVERSION%%/scheme/guile-procedures.txt
%%DATADIR%%/%%PORTVERSION%%/rewrite/checknas.rw
%%DATADIR%%/%%PORTVERSION%%/rewrite/log-hook.rw
%%DATADIR%%/%%PORTVERSION%%/rewrite/nas-ip.rw
%%MYSQL%%%%DATADIR%%/%%PORTVERSION%%/modules/mysql.so
%%PGSQL%%%%DATADIR%%/%%PORTVERSION%%/modules/postgres.so
%%EMACS%%%%EMACS_SITE_LISPDIR%%/radconf-mode.el
%%EMACS%%%%EMACS_SITE_LISPDIR%%/radius-mode.el
%%EMACS%%%%EMACS_SITE_LISPDIR%%/rewrite-mode.el
@dirrm etc/raddb/dict
@dirrm etc/raddb
@dirrm include/radius
@dirrm %%DATADIR%%/%%PORTVERSION%%/scheme
@dirrm %%DATADIR%%/%%PORTVERSION%%/rewrite
@dirrm %%DATADIR%%/%%PORTVERSION%%/modules
@dirrm %%DATADIR%%/%%PORTVERSION%%
@dirrm %%DATADIR%%
@unexec rmdir %D/lib/security 2>/dev/null || true
@unexec rm -f /var/run/radiusd.pid
@unexec rm -rf /var/log/radius
@unexec echo "Leaving /var/run/radutmp, /var/log/radwtmp and /var/log/radstat alone."
