# graphics/airsaned/files/airsaned.in
#!/bin/sh

# PROVIDE: airsaned
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf[.local] to enable airsaned:
#	airsaned_enable="YES"

. /etc/rc.subr

name=airsaned
rcvar=airsaned_enable

load_rc_config airsane

: ${airsaned_enable:="NO"}
: ${airsaned_args:=""}

command="%%PREFIX%%/sbin/${name}"
command_args="${airsaned_args} &"
sig_stop="KILL"

run_rc_command "$1"
