# www/tinyproxy/pkg-plist
etc/rc.d/tinyproxy.sh
etc/tinyproxy/tinyproxy.conf-dist
sbin/tinyproxy
share/tinyproxy/HTML_VARIABLES
share/tinyproxy/debug.html
share/tinyproxy/debug.html.dist
share/tinyproxy/default.html
share/tinyproxy/default.html.dist
share/tinyproxy/stats.html
share/tinyproxy/stats.html.dist
@dirrm share/tinyproxy
@dirrm etc/tinyproxy
