# net/sixxs-aiccu/pkg-plist
etc/aiccu.conf.sample
etc/rc.d/sixxs-aiccu.sh
sbin/sixxs-aiccu
%%EXAMPLESDIR%%/aiccu.conf
@dirrm %%EXAMPLESDIR%%
