# misc/compat4x/pkg-plist
etc/rc.d/000.compat4x.sh
%%CURRENT_ONLY:%%lib/compat/libc.so.4
%%CURRENT_ONLY:%%@unexec if [ -e %B/libc_r.so.4.compat3x ]; then ln -sf libc_r.so.4.compat3x %B/libc_r.so.4; else rm %B/libc_r.so.4; fi
%%CURRENT_ONLY:%%lib/compat/libc_r.so.4.compat4x
%%CURRENT_ONLY:%%@exec ln -sf %f %B/libc_r.so.4
%%FORBIDDEN:%%lib/compat/libcrypto.so.1
%%CURRENT_ONLY:%%lib/compat/libedit.so.3
lib/compat/libfetch.so.2
%%CURRENT_ONLY:%%@unexec if [ -e %B/libperl.so.3.compat3x ]; then ln -sf libperl.so.3.compat3x %B/libperl.so.3; else rm %B/libperl.so.3; fi
%%CURRENT_ONLY:%%lib/compat/libperl.so.3.compat4x
%%CURRENT_ONLY:%%@exec ln -sf %f %B/libperl.so.3
%%CURRENT_ONLY:%%lib/compat/libpam.so.1
%%CURRENT_ONLY:%%%%I386_ONLY:%%lib/compat/libposix1e.so.2
%%FORBIDDEN:%%lib/compat/libssl.so.1
%%CURRENT_ONLY:%%lib/compat/libstdc++.so.3
%%CURRENT_ONLY:%%lib/compat/libusb.so.0
%%CURRENT_ONLY:%%lib/compat/pam_cleartext_pass_ok.so
%%CURRENT_ONLY:%%lib/compat/pam_deny.so
%%CURRENT_ONLY:%%lib/compat/pam_opie.so
%%CURRENT_ONLY:%%lib/compat/pam_permit.so
%%CURRENT_ONLY:%%lib/compat/pam_radius.so
%%CURRENT_ONLY:%%lib/compat/pam_skey.so
%%CURRENT_ONLY:%%lib/compat/pam_ssh.so
%%CURRENT_ONLY:%%lib/compat/pam_tacplus.so
%%CURRENT_ONLY:%%lib/compat/pam_unix.so
@unexec rmdir %D/lib/compat 2>/dev/null || true
