# misc/compat4x/pkg-plist.freebsd4
etc/rc.d/000.compat4x.sh
%%FORBIDDEN:%%lib/compat/libcrypto.so.1
%%FORBIDDEN:%%lib/compat/libcrypto.so.2
lib/compat/libfetch.so.2
%%FORBIDDEN:%%lib/compat/libssl.so.1
%%FORBIDDEN:%%lib/compat/libssl.so.2
@unexec rmdir %D/lib/compat 2>/dev/null || true
