# net-p2p/torrentflux/pkg-plist
@comment $FreeBSD$
@unexec if cmp -s %D/etc/tfconfig.php %D/etc/tfconfig.php.default; then rm -f %D/etc/tfconfig.php; fi
etc/tfconfig.php.default
@exec [ -f %B/tfconfig.php ] || cp %B/%f %B/tfconfig.php
etc/rc.d/torrentflux.sh
libexec/btphptornado.py
%%DATADIR%%/torrentflux.sqlite
@dirrmtry %%DATADIR%%
%%PORTDOCS%%%%DOCSDIR%%/install.txt
%%PORTDOCS%%@dirrm %%DOCSDIR%%
