# net/openmdns/files/mdnsd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: mdnsd
# REQUIRE: NETWORKING
# BEFORE:  DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable this service:
#
# mdnsd_enable (bool):	Set it to YES to enable mdnsd(8) on startup.
#			Default: NO
# mdnsd_flags (string):	Arguments passed to mdnsd(8) such as the list of
#			network interfaces to listen on e.g., "em0 fxp0".
#			Default: lo0 (no network exposure)

. /etc/rc.subr

name="mdnsd"
desc="%%COMMENT%%"
rcvar="${name}_enable"
command="%%PREFIX%%/sbin/${name}"
extra_commands="reload"

: ${mdnsd_enable:="NO"}
: ${mdnsd_flags:="lo0"}

load_rc_config $name
run_rc_command "$1"
