# dns/knot/files/knotd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: knot
# REQUIRE: SERVERS cleanvar
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable knot:
#
# knot_enable="YES":	Set to NO by default.
#			Set it to YES to enable knot.
# knot_config="":	Set to %%PREFIX%%/etc/knot/knot.conf
#			by default.

. /etc/rc.subr

name=knot
rcvar=knot_enable

extra_commands="reload"
reload_cmd="${name}_x reload"
status_cmd="${name}_x status"
stop_cmd="  ${name}_x stop"

load_rc_config ${name}

: ${knot_enable:="NO"}
: ${knot_config="%%PREFIX%%/etc/knot/knot.conf"}

command=%%PREFIX%%/sbin/${name}d
command_args="-c ${knot_config} -d"
required_files=${knot_config}

knot_x() {
	%%PREFIX%%/sbin/${name}c -c ${knot_config} "$1"
	}

run_rc_command "$1"
