# net-mgmt/argus-monitor/files/argusd.in
#!/bin/sh

# PROVIDE: argusd
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# argusd_enable (bool):   Set to NO by default.
#               Set it to YES to enable argusd
#

. /etc/rc.subr

name="argusd"
rcvar=${name}_enable

command=%%PREFIX%%/sbin/${name}
# pidfile=/var/run/${name}.pid
required_files="
	%%PREFIX%%/etc/argus-monitor/log \
	%%PREFIX%%/etc/argus-monitor/notno \
	%%PREFIX%%/etc/argus-monitor/config \
	%%PREFIX%%/etc/argus-monitor/users"

stop_cmd="%%PREFIX%%/sbin/argusctl shutdown reason=\"rc.argusd stop\""

start_precmd="/usr/bin/touch %%PREFIX%%/etc/argus-monitor/log"
extra_commands="console hup"
console_cmd="%%PREFIX%%/sbin/argusctl -k console"
hup_cmd="%%PREFIX%%/sbin/argusctl -k hup"


load_rc_config $name

: ${argusd_enable="NO"}

run_rc_command "$1"
