# net/DarwinStreamingServer/files/streamingadminserver.sh.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: streamingadminserver
# REQUIRE: login
#
streamingadminserver_enable=${streamingadminserver_enable-"NO"}

. /etc/rc.subr

name=streamingadminserver
rcvar=streamingadminserver_enable
command=%%PREFIX%%/sbin/${name}.pl
command_interpreter=%%PERL%%

load_rc_config ${name}

pidfile=/var/run/streamingadminserver.pid

run_rc_command "$1"
