# print/acroreadwrapper/files/linux_adobe.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: linux_adobe
# REQUIRE: LOGIN abi
# KEYWORD: nojail

. /etc/rc.subr

name="linux_adobe"
rcvar=linux_adobe_enable
start_cmd="linux_adobe_start"
stop_cmd="linux_adobe_stop"

linux_adobe_start()
{
	load_kld %%PREFIX%%/%%KMODDIR%%/linux_adobe.ko
}
linux_adobe_stop()
{
	/sbin/kldunload linux_adobe
}

: ${linux_adobe_enable="YES"}
load_rc_config $name
run_rc_command "$1"
