# mail/mimedefang/pkg-plist
bin/mimedefang-multiplexor
bin/mimedefang
bin/mimedefang.pl
bin/mimedefang-util
bin/md-mx-ctrl
bin/watch-mimedefang
bin/watch-multiple-mimedefangs.tcl
etc/rc.d/mimedefang
man/man1/mimedefang-util.1.gz
man/man5/mimedefang-filter.5.gz
man/man7/mimedefang-notify.7.gz
man/man7/mimedefang-protocol.7.gz
man/man8/md-mx-ctrl.8.gz
man/man8/mimedefang-multiplexor.8.gz
man/man8/mimedefang.8.gz
man/man8/mimedefang.pl.8.gz
man/man8/watch-mimedefang.8.gz
man/man8/watch-multiple-mimedefangs.8.gz
@sample %%ETCDIR%%/mimedefang-ip-key.sample
@sample %%ETCDIR%%/mimedefang-filter.sample
@sample %%ETCDIR%%/sa-mimedefang.cf.sample
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README.SPAMASSASSIN
%%PORTDOCS%%%%DOCSDIR%%/README.SOPHIE
%%PORTDOCS%%%%DOCSDIR%%/README.NONROOT
%%PORTDOCS%%%%DOCSDIR%%/README.SECURITY
%%PORTDOCS%%%%DOCSDIR%%/README.TRU64
%%PORTDOCS%%%%DOCSDIR%%/README.IRIX
%%PORTDOCS%%%%DOCSDIR%%/README.VEXIRA
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/README
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/stream-by-domain-filter
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/suggested-minimum-filter-for-windows-clients
@dir(%%USER%%,%%GROUP%%,710) %%MIMEDEFANG_RUNDIR%%
@dir(%%USER%%,%%GROUP%%,700) %%MIMEDEFANG_QUARANTINEDIR%%
