# comms/qpage/pkg-plist
bin/qpage
etc/example-qpage.cf
etc/rc.d/qpage.sh
