# sysutils/moosefs2-master/files/openrc-mfsmetalogger.in
#!/sbin/openrc-run

name=mfsmetalogger

config_file=%%PREFIX%%/etc/mfs/${name}.cfg
required_files=${config_file}

command=%%PREFIX%%/sbin/${name}
pidfile=%%MFS_WORKDIR%%/mfs/.${name}.lock
command_args="-c $config_file"

depend() {
	keyword -shutdown
}
