# mail/missey/files/openrc-mps.in
#!/sbin/openrc-run

name="mps"
command="%%PREFIX%%/sbin/mps"
