# net/ushare/files/patch-ushare
--- scripts/ushare.orig	2008-10-13 13:23:21.000000000 +0800
+++ scripts/ushare	2008-10-13 13:24:47.000000000 +0800
@@ -1,82 +1,29 @@
-#!/bin/sh -e
+#!/bin/sh
 #
-# uShare init script
+# ushare.sh for rc.d usage (c) 2007 Volker Theile <votdev@gmx.de>.
 #
-### BEGIN INIT INFO
-# Provides:          ushare
-# Required-Start:    $local_fs $syslog $network
-# Should-Start:
-# Required-Stop:
-# Should-Stop:
-# Default-Start:     2 3 4 5
-# Default-Stop:      0 1 6
-# Short-Description: uShare
-# Description:       uShare UPnP (TM) A/V & DLNA Media Server
-#                    You should edit configuration in /etc/ushare.conf file
-#                    See http://ushare.geexbox.org for details
-### END INIT INFO
-
-PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
-DAEMON=/usr/bin/ushare
-NAME=ushare
-DESC="uShare UPnP A/V & DLNA Media Server"
-PIDFILE=/var/run/ushare.pid
-CONFIGFILE=/etc/ushare.conf
-
-# abort if no executable exists
-[ -x $DAEMON ] || exit 0
-
-# Get lsb functions
-. /lib/lsb/init-functions
-. /etc/default/rcS
-
-[ -f /etc/default/ushare ] && . /etc/default/ushare
-
-checkpid() {
-  [ -e $PIDFILE ] || touch $PIDFILE
-}
-
-check_shares() {
-  if [ -r "$CONFIGFILE" ]; then
-    . $CONFIGFILE
-    [ -n "$USHARE_DIR" ] && return 0
-  fi
-  return 1
-}
-
-case "$1" in
-  start)
-    log_daemon_msg "Starting $DESC: $NAME"
-    if ! $(check_shares); then
-      log_warning_msg "No shares avalaible ..."
-      log_end_msg 0
-    else
-      checkpid
-      start-stop-daemon --start --quiet --background --oknodo \
-        --make-pidfile --pidfile $PIDFILE \
-        --exec $DAEMON -- $USHARE_OPTIONS
-      log_end_msg $?
-    fi
-  ;;
-  stop)
-    log_daemon_msg "Stopping $DESC: $NAME"
-    start-stop-daemon --stop --signal 2 --quiet --oknodo --pidfile $PIDFILE
-    log_end_msg $?
-  ;;
-  reload|force-reload)
-    log_daemon_msg "Reloading $DESC: $NAME"
-    start-stop-daemon --stop --signal 1 --quiet --oknodo --pidfile $PIDFILE --exec $DAEMON
-    log_end_msg $?
-  ;;
-  restart)
-    $0 stop
-    $0 start
-  ;;
-  *)
-    N=/etc/init.d/$NAME
-    log_success_msg "Usage: $N {start|stop|restart|reload|force-reload}"
-    exit 1
-  ;;
-esac
 
-exit 0
+# PROVIDE: ushare
+# REQUIRE: DAEMON
+# BEFORE: LOGIN
+# KEYWORD: shutdown
+
+# Add the following line to /etc/rc.conf to enable `ushare':
+#
+#ushare_enable="YES"
+
+. /etc/rc.subr
+
+name=ushare
+rcvar=`set_rcvar`
+
+command="/usr/local/bin/$name"
+command_args="--daemon"
+required_files="/usr/local/etc/$name.conf"
+
+# read configuration and set defaults
+load_rc_config "$name"
+: ${ushare_enable="NO"}
+: ${ushare_flags=""}
+ 
+run_rc_command "$1"
