# cad/veroroute/pkg-plist
bin/veroroute
man/man1/veroroute.1.gz
share/applications/veroroute.desktop
share/gEDA/gafrc.d/veroroute-clib.scm
share/gEDA/sym/veroroute_battery/Battery.sym
share/gEDA/sym/veroroute_bbd/MN3006.sym
share/gEDA/sym/veroroute_bbd/MN3007.sym
share/gEDA/sym/veroroute_bbd/MN3009.sym
share/gEDA/sym/veroroute_bbd/MN3101.sym
share/gEDA/sym/veroroute_bbd/MN3102.sym
share/gEDA/sym/veroroute_bbd/MN3204.sym
share/gEDA/sym/veroroute_bbd/MN3206.sym
share/gEDA/sym/veroroute_bbd/MN3207.sym
share/gEDA/sym/veroroute_bbd/MN3209.sym
share/gEDA/sym/veroroute_bbd/SAD1024.sym
share/gEDA/sym/veroroute_connector/Pad.sym
share/gEDA/sym/veroroute_diode/Generic.sym
share/gEDA/sym/veroroute_diode/LED.sym
share/gEDA/sym/veroroute_diode/Schottky.sym
share/gEDA/sym/veroroute_diode/Zener.sym
share/gEDA/sym/veroroute_electromechanical/Potentiometer_Offboard.sym
share/gEDA/sym/veroroute_electromechanical/Potentiometer_Trimmer_Flat.sym
share/gEDA/sym/veroroute_electromechanical/Potentiometer_Trimmer_Vertical.sym
share/gEDA/sym/veroroute_electromechanical/Speaker_Offboard.sym
share/gEDA/sym/veroroute_electromechanical/Switch_DIP_3PST.sym
share/gEDA/sym/veroroute_electromechanical/Switch_DIP_DPST.sym
share/gEDA/sym/veroroute_electromechanical/Switch_OffBoard_DPDT.sym
share/gEDA/sym/veroroute_electromechanical/Switch_Offboard_3PDT.sym
share/gEDA/sym/veroroute_electromechanical/Switch_Offboard_DPST.sym
share/gEDA/sym/veroroute_electromechanical/Switch_Offboard_SPDT.sym
share/gEDA/sym/veroroute_electromechanical/Switch_Offboard_SPST.sym
share/gEDA/sym/veroroute_electromechanical/Switch_Onboard_3PDT.sym
share/gEDA/sym/veroroute_electromechanical/Switch_Onboard_DPDT.sym
share/gEDA/sym/veroroute_electromechanical/Switch_Onboard_DPST.sym
share/gEDA/sym/veroroute_electromechanical/Switch_Onboard_SPDT.sym
share/gEDA/sym/veroroute_electromechanical/Switch_Onboard_SPST.sym
share/gEDA/sym/veroroute_linear/Comparator_LM311.sym
share/gEDA/sym/veroroute_linear/Comparator_LM319.sym
share/gEDA/sym/veroroute_linear/Comparator_LM339.sym
share/gEDA/sym/veroroute_linear/Generic_DIP10.sym
share/gEDA/sym/veroroute_linear/Generic_DIP12.sym
share/gEDA/sym/veroroute_linear/Generic_DIP14.sym
share/gEDA/sym/veroroute_linear/Generic_DIP16.sym
share/gEDA/sym/veroroute_linear/Generic_DIP8.sym
share/gEDA/sym/veroroute_linear/Generic_SIP10.sym
share/gEDA/sym/veroroute_linear/Generic_SIP12.sym
share/gEDA/sym/veroroute_linear/Generic_SIP14.sym
share/gEDA/sym/veroroute_linear/Generic_SIP16.sym
share/gEDA/sym/veroroute_linear/Generic_SIP8.sym
share/gEDA/sym/veroroute_linear/OTA_LM13600.sym
share/gEDA/sym/veroroute_linear/OpAmp_LM324.sym
share/gEDA/sym/veroroute_linear/OpAmp_LM358.sym
share/gEDA/sym/veroroute_linear/OpAmp_LM741.sym
share/gEDA/sym/veroroute_linear/OpAmp_NJM4558L.sym
share/gEDA/sym/veroroute_linear/OpAmp_TL072.sym
share/gEDA/sym/veroroute_linear/OpAmp_TL074.sym
share/gEDA/sym/veroroute_linear/Timer_NE555.sym
share/gEDA/sym/veroroute_logic/24LC32A.sym
share/gEDA/sym/veroroute_logic/4013.sym
share/gEDA/sym/veroroute_logic/4049.sym
share/gEDA/sym/veroroute_logic/4050.sym
share/gEDA/sym/veroroute_passive/Capacitor_Ceramic.sym
share/gEDA/sym/veroroute_passive/Capacitor_Ceramic_100mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Ceramic_200mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Ceramic_300mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Electrolytic_200mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Electrolytic_250mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Electrolytic_300mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Electrolytic_400mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Electrolytic_500mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Electrolytic_600mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Electrolytic_NP_200mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Electrolytic_NP_250mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Electrolytic_NP_300mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Electrolytic_NP_400mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Electrolytic_NP_500mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Electrolytic_NP_600mil.sym
share/gEDA/sym/veroroute_passive/Capacitor_Film.sym
share/gEDA/sym/veroroute_passive/Capacitor_Film_100.sym
share/gEDA/sym/veroroute_passive/Capacitor_Film_200.sym
share/gEDA/sym/veroroute_passive/Capacitor_Film_300.sym
share/gEDA/sym/veroroute_passive/Crystal.sym
share/gEDA/sym/veroroute_passive/Resistor.sym
share/gEDA/sym/veroroute_passive/Resistor_100mil.sym
share/gEDA/sym/veroroute_passive/Resistor_200mil.sym
share/gEDA/sym/veroroute_passive/Resistor_300mil.sym
share/gEDA/sym/veroroute_passive/Resistor_400mil.sym
share/gEDA/sym/veroroute_regulator/LM137.sym
share/gEDA/sym/veroroute_regulator/LM317.sym
share/gEDA/sym/veroroute_regulator/LM336.sym
share/gEDA/sym/veroroute_regulator/LM337.sym
share/gEDA/sym/veroroute_regulator/LM385.sym
share/gEDA/sym/veroroute_regulator/LM78Lxx.sym
share/gEDA/sym/veroroute_regulator/LM78xx.sym
share/gEDA/sym/veroroute_regulator/LM79Lxx.sym
share/gEDA/sym/veroroute_regulator/LM79xx.sym
share/gEDA/sym/veroroute_regulator/LP2950.sym
share/gEDA/sym/veroroute_regulator/TL431.sym
share/gEDA/sym/veroroute_spin/Spin_FV1.sym
share/gEDA/sym/veroroute_transistor/Darlington_NPN_TIP100.sym
share/gEDA/sym/veroroute_transistor/Darlington_NPN_TIP101.sym
share/gEDA/sym/veroroute_transistor/Darlington_NPN_TIP102.sym
share/gEDA/sym/veroroute_transistor/Darlington_PNP_TIP105.sym
share/gEDA/sym/veroroute_transistor/Darlington_PNP_TIP106.sym
share/gEDA/sym/veroroute_transistor/Darlington_PNP_TIP107.sym
share/gEDA/sym/veroroute_transistor/JFET_Nchannel_2N5457.sym
share/gEDA/sym/veroroute_transistor/JFET_Nchannel_2N5458.sym
share/gEDA/sym/veroroute_transistor/JFET_Nchannel_2N5459.sym
share/gEDA/sym/veroroute_transistor/JFET_Nchannel_2N5484.sym
share/gEDA/sym/veroroute_transistor/JFET_Nchannel_2N5485.sym
share/gEDA/sym/veroroute_transistor/JFET_Nchannel_2N5486.sym
share/gEDA/sym/veroroute_transistor/JFET_Nchannel_2SK117.sym
share/gEDA/sym/veroroute_transistor/JFET_Nchannel_2SK30.sym
share/gEDA/sym/veroroute_transistor/JFET_Nchannel_BF245.sym
share/gEDA/sym/veroroute_transistor/JFET_Nchannel_J201.sym
share/gEDA/sym/veroroute_transistor/JFET_Nchannel_J202.sym
share/gEDA/sym/veroroute_transistor/JFET_Nchannel_MPF102.sym
share/gEDA/sym/veroroute_transistor/MOSFET_Nchannel_2N7000.sym
share/gEDA/sym/veroroute_transistor/MOSFET_Nchannel_BS170.sym
share/gEDA/sym/veroroute_transistor/MOSFET_Nchannel_MTP3055.sym
share/gEDA/sym/veroroute_transistor/MOSFET_Pchannel_MTP2955.sym
share/gEDA/sym/veroroute_transistor/NPN_2N2222.sym
share/gEDA/sym/veroroute_transistor/NPN_2N3903.sym
share/gEDA/sym/veroroute_transistor/NPN_2N3904.sym
share/gEDA/sym/veroroute_transistor/NPN_2N5088.sym
share/gEDA/sym/veroroute_transistor/NPN_2N5089.sym
share/gEDA/sym/veroroute_transistor/NPN_BC182.sym
share/gEDA/sym/veroroute_transistor/NPN_BC183.sym
share/gEDA/sym/veroroute_transistor/NPN_BC184.sym
share/gEDA/sym/veroroute_transistor/NPN_BC212.sym
share/gEDA/sym/veroroute_transistor/NPN_BC213.sym
share/gEDA/sym/veroroute_transistor/NPN_BC214.sym
share/gEDA/sym/veroroute_transistor/NPN_BC237.sym
share/gEDA/sym/veroroute_transistor/NPN_BC238.sym
share/gEDA/sym/veroroute_transistor/NPN_BC239.sym
share/gEDA/sym/veroroute_transistor/NPN_BC337.sym
share/gEDA/sym/veroroute_transistor/NPN_BC338.sym
share/gEDA/sym/veroroute_transistor/NPN_BC546.sym
share/gEDA/sym/veroroute_transistor/NPN_BC547.sym
share/gEDA/sym/veroroute_transistor/NPN_BC548.sym
share/gEDA/sym/veroroute_transistor/NPN_BC549.sym
share/gEDA/sym/veroroute_transistor/NPN_BC550.sym
share/gEDA/sym/veroroute_transistor/NPN_BC635.sym
share/gEDA/sym/veroroute_transistor/NPN_BC637.sym
share/gEDA/sym/veroroute_transistor/NPN_BC639.sym
share/gEDA/sym/veroroute_transistor/NPN_MJE3055.sym
share/gEDA/sym/veroroute_transistor/NPN_TIP31.sym
share/gEDA/sym/veroroute_transistor/NPN_TIP41.sym
share/gEDA/sym/veroroute_transistor/PNP_2N3905.sym
share/gEDA/sym/veroroute_transistor/PNP_2N3906.sym
share/gEDA/sym/veroroute_transistor/PNP_2N5086.sym
share/gEDA/sym/veroroute_transistor/PNP_2N5087.sym
share/gEDA/sym/veroroute_transistor/PNP_BC307.sym
share/gEDA/sym/veroroute_transistor/PNP_BC308.sym
share/gEDA/sym/veroroute_transistor/PNP_BC309.sym
share/gEDA/sym/veroroute_transistor/PNP_BC327.sym
share/gEDA/sym/veroroute_transistor/PNP_BC328.sym
share/gEDA/sym/veroroute_transistor/PNP_BC556.sym
share/gEDA/sym/veroroute_transistor/PNP_BC557.sym
share/gEDA/sym/veroroute_transistor/PNP_BC558.sym
share/gEDA/sym/veroroute_transistor/PNP_BC559.sym
share/gEDA/sym/veroroute_transistor/PNP_BC560.sym
share/gEDA/sym/veroroute_transistor/PNP_BC640.sym
share/gEDA/sym/veroroute_transistor/PNP_MJE2955.sym
share/gEDA/sym/veroroute_transistor/PNP_TIP32.sym
share/gEDA/sym/veroroute_transistor/PNP_TIP42.sym
share/pixmaps/veroroute.png
%%DATADIR%%/tutorials/tutorial_0.vrt
%%DATADIR%%/tutorials/tutorial_1.vrt
%%DATADIR%%/tutorials/tutorial_10.vrt
%%DATADIR%%/tutorials/tutorial_11.vrt
%%DATADIR%%/tutorials/tutorial_12.vrt
%%DATADIR%%/tutorials/tutorial_13.vrt
%%DATADIR%%/tutorials/tutorial_14.vrt
%%DATADIR%%/tutorials/tutorial_15.vrt
%%DATADIR%%/tutorials/tutorial_16.vrt
%%DATADIR%%/tutorials/tutorial_17.vrt
%%DATADIR%%/tutorials/tutorial_18.vrt
%%DATADIR%%/tutorials/tutorial_19.vrt
%%DATADIR%%/tutorials/tutorial_2.vrt
%%DATADIR%%/tutorials/tutorial_20.vrt
%%DATADIR%%/tutorials/tutorial_21.vrt
%%DATADIR%%/tutorials/tutorial_22.vrt
%%DATADIR%%/tutorials/tutorial_3.vrt
%%DATADIR%%/tutorials/tutorial_4.vrt
%%DATADIR%%/tutorials/tutorial_5.vrt
%%DATADIR%%/tutorials/tutorial_6.vrt
%%DATADIR%%/tutorials/tutorial_7.vrt
%%DATADIR%%/tutorials/tutorial_8.vrt
%%DATADIR%%/tutorials/tutorial_9.vrt
%%DATADIR%%/veroroute.png
