# net/ndisc6/files/openrc-rdnssd.in
#!/sbin/openrc-run

name=rdnssd

pidfile=/var/run/${name}.pid
command="%%PREFIX%%/sbin/${name}"
command_args="-p ${pidfile}"

depend() {
	use net
}
