# net-mgmt/hawk/files/openrc-hawk.in
#!/sbin/openrc-run

name=hawk

command=%%PREFIX%%/sbin/hawk
command_args='&'

depend() {
	need mysql
	keyword -shutdown
}
