# ftp/proftpd/pkg-plist
@unexec if cmp -s %D/etc/proftpd.conf %D/etc/proftpd.conf.default; then rm -f %D/etc/proftpd.conf; fi
bin/ftpcount
bin/ftpwho
etc/proftpd.conf.default
etc/rc.d/proftpd.sh.sample
@exec if [ ! -f %D/etc/proftpd.conf ] ; then cp -p %D/%F %B/proftpd.conf; fi
libexec/proftpd
sbin/ftpshut
share/doc/proftpd/Configuration.html
share/doc/proftpd/faq.html
@dirrm share/doc/proftpd
