# devel/distcc/pkg-plist
@comment $FreeBSD$
bin/distcc
bin/distccmon-text
%%WITH_GNOME%%bin/distccmon-gnome
sbin/distccd
etc/rc.d/distccd.sh
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/protocol-1.txt
%%PORTDOCS%%%%DOCSDIR%%/protocol-2.txt
%%PORTDOCS%%%%DOCSDIR%%/status-1.txt
%%PORTDOCS%%%%DOCSDIR%%/reporting-bugs.txt
%%PORTDOCS%%@dirrm %%DOCSDIR%%
