# sysutils/bmd/pkg-plist
etc/rc.d/bmd
etc/bmd.conf.example
include/bmd_plugin.h
sbin/bmd
sbin/bmdctl
share/man/man5/bmd.conf.5.gz
share/man/man8/bmd.8.gz
share/man/man8/bmdctl.8.gz
@dir libexec/bmd
@dir var/cache/bmd
