# net/yate/files/yate.sh.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: yate
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable the YATE:
#
# yate_enable="YES"
#

. /etc/rc.subr

name=yate
rcvar=yate_enable

command="%%PREFIX%%/bin/yate"
pidfile="/var/run/${name}.pid"
logfile="/var/log/yate/messages"
cdrfile=`awk -F= '/^file=.+/ {print $2}' %%PREFIX%%/etc/yate/cdrfile.conf`

start_precmd="yate_prestart"
stop_postcmd="yate_poststop"

yate_prestart()
{
    mkdir -p $(dirname $pidfile)
    chown yate:yate $(dirname $pidfile)
    mkdir -p $(dirname $cdrfile)
    mkdir -p $(dirname $logfile)
    chown yate:yate $(dirname $cdrfile)
    chown yate:yate $(dirname $logfile)
}

yate_poststop()
{
    rm -f $pidfile
}

yate_enable=${yate_enable="NO"}
yate_flags=${yate_flags-"-Du -d -l ${logfile} -p ${pidfile} -s"}

load_rc_config ${name}
run_rc_command "$1"
