# security/dropbear/pkg-plist
bin/dbclient
bin/dropbearkey
bin/dropbearconvert
etc/rc.d/dropbear.sh.sample
sbin/dropbear
@unexec rmdir %D/etc/dropbear 2>/dev/null || echo "===> You can safely remove %D/etc/dropbear directory, if you don't want to use dropbear anymore."
