# audio/rplay/pkg-plist
etc/rc.d/rplayd.sh
etc/rplay.hosts
bin/rplay
bin/rptp
@comment requires XForms
@comment bin/xrplay
sbin/rplayd
include/rplay.h
lib/librplay.a
@unexec install-info --delete
info/librplay.info
info/rplay.info
info/RPLAY.info
info/RPTP.info
@exec install-info
