# mail/rspamd-devel/files/rspamd_redirector.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: rspamd_redirector
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable rspamd_redirector:
#
# rspamd_redirector_enable (bool):	Set to "NO" by default.
# 					Set it to "YES" to enable redirector.

. /etc/rc.subr

name="rspamd_redirector"
rcvar=rspamd_redirector_enable

load_rc_config $name

: ${rspamd_redirector_enable:="NO"}

pidfile=${rspamd_redirector_pidfile:-"/var/run/rspamd/redirector.pid"}

command="%%PREFIX%%/bin/rspamd-redirector"
command_interpreter="%%PREFIX%%/bin/perl"
extra_commands="reload"
sig_reload="USR1"
stop_postcmd="rm -f $pidfile"

run_rc_command "$1"
