# sysutils/firstboot-pkgs/files/openrc-firstboot_pkgs.in
#!/sbin/openrc-run

# and place a list of packages in firstboot_pkgs_list, e.g.,
#
# firstboot_pkgs_list="apache22 php5 mysql56-server"

name="firstboot_pkgs"
command="pkg"

start() {

	# Count rc.d scripts
	nscriptso=`ls /usr/local/etc/rc.d | wc -l`

	# Bootstrap if necessary
	if ! pkg -N 2>/dev/null; then
		env ASSUME_ALWAYS_YES=YES pkg bootstrap |
		    cat
	fi

	# Install requested packages, if any
	for package in ${firstboot_pkgs_list}; do
		env ASSUME_ALWAYS_YES=YES pkg install ${package} </dev/null |
		    cat
	done

	# Count rc.d scripts again
	nscriptsn=`ls /usr/local/etc/rc.d | wc -l`

	# If we have more scripts, request a reboot
	if [ $nscriptsn -ne $nscriptso ]; then
		einfo "Requesting reboot after installing packages with rc.d scripts."
		touch ${firstboot_sentinel}-reboot
	fi
}

stop() {
}

depend() {
	use net
}
