# sysutils/afbinit/files/afbinit.sh.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: afbinit
# REQUIRE: devfs
#
# Available configuration variables for afbinit are:
#
# afbinit_enable (bool):	Set to "YES" to enable afbinit.
#				Defaults to "NO".
#
# Add at least the following line to /etc/rc.conf or /etc/rc.conf.local to
# enable afbinit:
#
# afbinit="YES"
#

. /etc/rc.subr

name="afbinit"
rcvar=afbinit_enable

start_cmd="afbinit_startcmd"

afbinit_startcmd()
{
	if [ ! -r %%DATADIR%%/afb.ucode ]; then
		err 1 '%%DATADIR%%/afb.ucode is not readable.'
	fi

	for fb in /dev/fb[0-7]; do
		%%PREFIX%%/sbin/afbinit ${fb} %%DATADIR%%/afb.ucode
	done
}

load_rc_config $name

: ${afbinit_enable="NO"}

run_rc_command "$1"
