# mail/up-imapproxy/pkg-plist
sbin/pimpstat
sbin/in.imapproxyd
etc/rc.d/imapproxyd.sh
etc/imapproxyd.conf.default
etc/imapproxyd.conf.sample
