# net/isc-dhcp43-relay/pkg-plist
%%IPV6%%etc/rc.d/isc-dhcrelay6
sbin/dhcrelay
man/man8/dhcrelay.8.gz
