# mail/spamd/pkg-plist
libexec/spamd
sbin/spamd-setup
etc/rc.d/spamd.sh
etc/spamd.conf.sample
