# mail/spamd/pkg-plist
libexec/spamd
libexec/spamlogd
sbin/spamd-setup
sbin/spamdb
etc/rc.d/pfspamd.sh
etc/spamd.conf.sample
