# dns/totd/files/openrc-totd.in
#!/sbin/openrc-run

name="totd"
command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="%%PREFIX%%/etc/${name}.conf"
