# x11/pcdm/pkg-plist
bin/PCDM-session
bin/pc-dmconf
bin/start-pcdm-pico
@sample etc/pcdm.conf.dist etc/pcdm.conf
etc/rc.d/pcdm
etc/sudoers.d/pico-sudoers
sbin/PCDMd
share/PCDM/i18n/PCDM_af.qm
share/PCDM/i18n/PCDM_ar.qm
share/PCDM/i18n/PCDM_az.qm
share/PCDM/i18n/PCDM_bg.qm
share/PCDM/i18n/PCDM_bn.qm
share/PCDM/i18n/PCDM_bs.qm
share/PCDM/i18n/PCDM_ca.qm
share/PCDM/i18n/PCDM_cs.qm
share/PCDM/i18n/PCDM_cy.qm
share/PCDM/i18n/PCDM_da.qm
share/PCDM/i18n/PCDM_de.qm
share/PCDM/i18n/PCDM_el.qm
share/PCDM/i18n/PCDM_en.qm
share/PCDM/i18n/PCDM_en_GB.qm
share/PCDM/i18n/PCDM_en_AU.qm
share/PCDM/i18n/PCDM_en_US.qm
share/PCDM/i18n/PCDM_en_ZA.qm
share/PCDM/i18n/PCDM_es.qm
share/PCDM/i18n/PCDM_et.qm
share/PCDM/i18n/PCDM_eu.qm
share/PCDM/i18n/PCDM_fa.qm
share/PCDM/i18n/PCDM_fi.qm
share/PCDM/i18n/PCDM_fr.qm
share/PCDM/i18n/PCDM_fr_CA.qm
share/PCDM/i18n/PCDM_fur.qm
share/PCDM/i18n/PCDM_gl.qm
share/PCDM/i18n/PCDM_he.qm
share/PCDM/i18n/PCDM_hi.qm
share/PCDM/i18n/PCDM_hr.qm
share/PCDM/i18n/PCDM_hu.qm
share/PCDM/i18n/PCDM_id.qm
share/PCDM/i18n/PCDM_is.qm
share/PCDM/i18n/PCDM_it.qm
share/PCDM/i18n/PCDM_ja.qm
share/PCDM/i18n/PCDM_ka.qm
share/PCDM/i18n/PCDM_ko.qm
share/PCDM/i18n/PCDM_lt.qm
share/PCDM/i18n/PCDM_lv.qm
share/PCDM/i18n/PCDM_mk.qm
share/PCDM/i18n/PCDM_mn.qm
share/PCDM/i18n/PCDM_ms.qm
share/PCDM/i18n/PCDM_mt.qm
share/PCDM/i18n/PCDM_nb.qm
share/PCDM/i18n/PCDM_ne.qm
share/PCDM/i18n/PCDM_nl.qm
share/PCDM/i18n/PCDM_pa.qm
share/PCDM/i18n/PCDM_pl.qm
share/PCDM/i18n/PCDM_pt.qm
share/PCDM/i18n/PCDM_pt_BR.qm
share/PCDM/i18n/PCDM_ro.qm
share/PCDM/i18n/PCDM_ru.qm
share/PCDM/i18n/PCDM_sa.qm
share/PCDM/i18n/PCDM_sk.qm
share/PCDM/i18n/PCDM_sl.qm
share/PCDM/i18n/PCDM_sr.qm
share/PCDM/i18n/PCDM_sv.qm
share/PCDM/i18n/PCDM_sw.qm
share/PCDM/i18n/PCDM_ta.qm
share/PCDM/i18n/PCDM_tg.qm
share/PCDM/i18n/PCDM_th.qm
share/PCDM/i18n/PCDM_tr.qm
share/PCDM/i18n/PCDM_uk.qm
share/PCDM/i18n/PCDM_ur.qm
share/PCDM/i18n/PCDM_uz.qm
share/PCDM/i18n/PCDM_vi.qm
share/PCDM/i18n/PCDM_zh_CN.qm
share/PCDM/i18n/PCDM_zh_HK.qm
share/PCDM/i18n/PCDM_zh_TW.qm
share/PCDM/i18n/PCDM_zu.qm
share/PCDM/pcdm-session
share/PCDM/themes/default/banner.png
share/PCDM/themes/default/default-desktop.png
share/PCDM/themes/default/default-user.png
share/PCDM/themes/default/default.theme
share/PCDM/themes/default/keyboardlayout.png
share/PCDM/themes/default/localechange.png
share/PCDM/themes/default/login-stealth.png
share/PCDM/themes/default/login.png
share/PCDM/themes/default/nextDE.png
share/PCDM/themes/default/password.png
share/PCDM/themes/default/pcdm-default-4k.jpg
share/PCDM/themes/default/previousDE.png
share/PCDM/themes/default/system.png
share/PCDM/themes/default/virtualkeyboard.png
share/PCDM/themes/dragonflybsd/default-user.png
share/PCDM/themes/dragonflybsd/dragonflybsd.theme
share/PCDM/themes/dragonflybsd/keyboardlayout.png
share/PCDM/themes/dragonflybsd/localechange.png
share/PCDM/themes/dragonflybsd/login-stealth.png
share/PCDM/themes/dragonflybsd/login.png
share/PCDM/themes/dragonflybsd/nextDE.png
share/PCDM/themes/dragonflybsd/password.png
share/PCDM/themes/dragonflybsd/pcdm-dragonflybsd-4k.jpg
share/PCDM/themes/dragonflybsd/previousDE.png
share/PCDM/themes/dragonflybsd/system.png
share/PCDM/themes/dragonflybsd/virtualkeyboard.png
share/PCDM/themes/enterprise/banner.png
share/PCDM/themes/enterprise/default-desktop.png
share/PCDM/themes/enterprise/default-user.png
share/PCDM/themes/enterprise/enterprise.theme
share/PCDM/themes/enterprise/keyboardlayout.png
share/PCDM/themes/enterprise/localechange.png
share/PCDM/themes/enterprise/login-stealth.png
share/PCDM/themes/enterprise/login.png
share/PCDM/themes/enterprise/nextDE.png
share/PCDM/themes/enterprise/password.png
share/PCDM/themes/enterprise/pcdm-enterprise-4k.jpg
share/PCDM/themes/enterprise/previousDE.png
share/PCDM/themes/enterprise/system.png
share/PCDM/themes/enterprise/virtualkeyboard.png
share/PCDM/themes/freebsd/default-desktop.png
share/PCDM/themes/freebsd/default-user.png
share/PCDM/themes/freebsd/freebsd.theme
share/PCDM/themes/freebsd/keyboardlayout.png
share/PCDM/themes/freebsd/localechange.png
share/PCDM/themes/freebsd/login-stealth.png
share/PCDM/themes/freebsd/login.png
share/PCDM/themes/freebsd/nextDE.png
share/PCDM/themes/freebsd/password.png
share/PCDM/themes/freebsd/pcdm-freebsd-4k.jpg
share/PCDM/themes/freebsd/previousDE.png
share/PCDM/themes/freebsd/system.png
share/PCDM/themes/freebsd/virtualkeyboard.png
share/PCDM/themes/hardenedbsd/default-desktop.png
share/PCDM/themes/hardenedbsd/default-user.png
share/PCDM/themes/hardenedbsd/hardenedbsd.theme
share/PCDM/themes/hardenedbsd/keyboardlayout.png
share/PCDM/themes/hardenedbsd/localechange.png
share/PCDM/themes/hardenedbsd/login-stealth.png
share/PCDM/themes/hardenedbsd/login.png
share/PCDM/themes/hardenedbsd/nextDE.png
share/PCDM/themes/hardenedbsd/password.png
share/PCDM/themes/hardenedbsd/pcdm-hardenedbsd-4k.jpg
share/PCDM/themes/hardenedbsd/previousDE.png
share/PCDM/themes/hardenedbsd/system.png
share/PCDM/themes/hardenedbsd/virtualkeyboard.png
share/PCDM/themes/netbsd/default-user.png
share/PCDM/themes/netbsd/keyboardlayout.png
share/PCDM/themes/netbsd/localechange.png
share/PCDM/themes/netbsd/login-stealth.png
share/PCDM/themes/netbsd/login.png
share/PCDM/themes/netbsd/netbsd.theme
share/PCDM/themes/netbsd/nextDE.png
share/PCDM/themes/netbsd/password.png
share/PCDM/themes/netbsd/pcdm-netbsd-4k.jpg
share/PCDM/themes/netbsd/previousDE.png
share/PCDM/themes/netbsd/system.png
share/PCDM/themes/netbsd/virtualkeyboard.png
share/PCDM/themes/openbsd/default-user.png
share/PCDM/themes/openbsd/keyboardlayout.png
share/PCDM/themes/openbsd/localechange.png
share/PCDM/themes/openbsd/login-stealth.png
share/PCDM/themes/openbsd/login.png
share/PCDM/themes/openbsd/nextDE.png
share/PCDM/themes/openbsd/openbsd.theme
share/PCDM/themes/openbsd/password.png
share/PCDM/themes/openbsd/pcdm-openbsd-4k.jpg
share/PCDM/themes/openbsd/previousDE.png
share/PCDM/themes/openbsd/system.png
share/PCDM/themes/openbsd/virtualkeyboard.png
share/PCDM/themes/trueos/banner.png
share/PCDM/themes/trueos/default-user.png
share/PCDM/themes/trueos/keyboardlayout.png
share/PCDM/themes/trueos/localechange.png
share/PCDM/themes/trueos/login-stealth.png
share/PCDM/themes/trueos/login.png
share/PCDM/themes/trueos/nextDE.png
share/PCDM/themes/trueos/password.png
share/PCDM/themes/trueos/pcdm-trueos-4k.jpg
share/PCDM/themes/trueos/previousDE.png
share/PCDM/themes/trueos/system.png
share/PCDM/themes/trueos/trueos.theme
share/PCDM/themes/trueos/virtualkeyboard.png
share/applications/pc-dmconf.desktop
share/pixmaps/pc-dmconf.png
share/xsessions/_deforaos.desktop
share/xsessions/_evilwm.desktop
share/xsessions/_fluxbox.desktop
share/xsessions/_fvwm-crystal.desktop
share/xsessions/_fvwm.desktop
share/xsessions/_gnome-classic.desktop
share/xsessions/_icewm.desktop
share/xsessions/_kde4.desktop
share/xsessions/_plex.desktop
share/xsessions/_ratpoison.desktop
share/xsessions/_spectrwm.desktop
share/xsessions/_vtwm.desktop
share/xsessions/_windowlab.desktop
share/xsessions/_windowmaker.desktop
