# sysutils/nagios-statd/pkg-plist
@comment $FreeBSD$
bin/nagios-stat
sbin/nagios-statd
etc/rc.d/nagios-statd.sh.sample
