# net-mgmt/nrpe2/pkg-plist
sbin/nrpe2
etc/nrpe.cfg-sample
etc/rc.d/nrpe.sh
libexec/nagios/check_nrpe2
