# net/zephyr/pkg-plist
bin/xzwrite
bin/zaway
bin/zctl
bin/zleave
bin/zlocate
bin/zmailnotify
bin/znol
bin/zstat
bin/zwgc
bin/zwrite
etc/rc.d/zhm.sh
etc/zephyr/default.subscriptions
include/zephyr/mit-copyright.h
include/zephyr/zephyr.h
include/zephyr/zephyr_err.h
lib/libzephyr.a
sbin/zephyrd
sbin/zhm
sbin/zpopnotify
sbin/zshutdown_notify
%%PORTDOCS%%share/doc/zephyr/INSTALL
%%PORTDOCS%%share/doc/zephyr/NOTES
%%PORTDOCS%%share/doc/zephyr/OPERATING
%%PORTDOCS%%share/doc/zephyr/README
%%PORTDOCS%%share/doc/zephyr/USING
share/zephyr/XZwrite
share/zephyr/xzwrite.bitmap
share/zephyr/zwgc.desc
share/zephyr/zwgc_resources
@dirrm etc/zephyr
@dirrm include/zephyr
%%PORTDOCS%%@dirrm share/doc/zephyr
@dirrm share/zephyr
