# security/ipguard/pkg-plist
sbin/ipguard
etc/rc.d/ipguard.sh.sample
