# net/scribe/files/scribe.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: scribe
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name="scribe"
rcvar=`set_rcvar`
pidfile="/var/run/${name}.pid"
command="%%PREFIX%%/bin/scribed"

start_cmd=scribe_start
stop_postcmd=scribe_cleanup

scribe_start()
{
    echo "Starting ${name}."
    /usr/sbin/daemon -cf -p ${pidfile} ${command} ${scribe_flags} 2>>/var/log/scribe.log
}

scribe_cleanup() {
    [ -f ${pidfile} ] && /bin/unlink ${pidfile}
}

load_rc_config "$name"

: ${scribe_enable="NO"}

run_rc_command "$1"
