# misc/trueos-desktop/pkg-plist
etc/sudoers.d
etc/sudoers.d/trueos-desktop
etc/init.d/trueosinit
etc/rc.d/trueosinit
etc/trueos.conf.dist
etc/trueos/recommended-services-desktop
etc/trueos/required-services-desktop
share/trueos/xstartup-parallel/disable-beep.sh
share/trueos/xstartup-parallel/load-xresources.sh
share/trueos/xstartup/enable-inputmethod.sh
share/trueos/xstartup/ssh-agent.sh
share/trueos/xstartup/02_fix-gtk-environ.sh
share/trueos/xstartup/01_xdg-user-dirs.sh
share/trueos/xstartup/gpg-agent.sh
share/wallpapers/PCBSD/b_1280x1024.jpg
share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1440x900_v2.jpg
share/wallpapers/PCBSD/g_1280x1024.jpg
share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1280x1024_v2.jpg
share/wallpapers/PCBSD/PCBSD_10_Joule_WP_600x400_v2.jpg
share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1920x1080_v2.jpg
share/wallpapers/PCBSD/PCBSD_Enterprise_WP_1920x1080.jpg
share/wallpapers/PCBSD/e_1280x1024.jpg
share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1280x800_v2.jpg
share/wallpapers/PCBSD/a_1280x1024.jpg
share/wallpapers/PCBSD/d_1280x1024.jpg
share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1920x1200_v2.jpg
share/wallpapers/PCBSD/a2_1280x1024.jpg
share/wallpapers/PCBSD/f2_1280x1024.jpg
share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1600x1200_v2.jpg
share/wallpapers/PCBSD/PCBSD_10_Joule_WP_800x600_v2.jpg
share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1024x768_v2.jpg
share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1024x600_v2.jpg
share/wallpapers/PCBSD/f_1280x1024.jpg
share/wallpapers/PCBSD/c_1280x1024.jpg
share/wallpapers/FreeBSD/d_1280x1024.jpg
share/wallpapers/FreeBSD/a_1280x1024.jpg
share/wallpapers/FreeBSD/a_1600x900.jpg
share/wallpapers/FreeBSD/c_1600x900.jpg
share/wallpapers/FreeBSD/b_1600x900.jpg
share/wallpapers/FreeBSD/g_1600x900.jpg
share/wallpapers/FreeBSD/c_1280x1024.jpg
share/wallpapers/FreeBSD/f_1600x900.jpg
share/wallpapers/FreeBSD/d_1600x900.jpg
share/wallpapers/FreeBSD/e_1600x900.jpg
share/wallpapers/FreeBSD/f_1280x1024.jpg
share/wallpapers/FreeBSD/g_1280x1024.jpg
share/wallpapers/FreeBSD/b_1280x1024.jpg
share/wallpapers/FreeBSD/e_1280x1024.jpg
share/wallpapers/TrueOS/trueos-2-1080.png
share/wallpapers/TrueOS/trueos-2-4k.png
share/wallpapers/TrueOS/trueos-current-underconstruction-1080.png
share/wallpapers/TrueOS/TrueOS_17-12_4K_grey.jpg
share/wallpapers/TrueOS/trueos-current-underconstruction-4k.png
share/wallpapers/TrueOS/pcdm-trueos-1080.png
share/wallpapers/TrueOS/pcdm-trueos-4k.png
share/wallpapers/TrueOS/trueos-1-4k.png
share/wallpapers/TrueOS/trueos-1-1080.png
share/wallpapers/TrueOS/TrueOS_17-12_4K_purple.jpg
share/wallpapers/trueos-default.jpg
share/icons/hicolor/256x256/places/start-here-lumina.png
share/icons/default/index.theme
share/cups/data/default-testpage.orig.pdf
share/fluxbox/styles/bora_blue_4k/theme.cfg
/usr/share/skel/dot.xprofile
@sample etc/trueos.conf etc/trueos.conf.dist
@sample /etc/devfs.rules.trueos /etc/devfs.rules
