# devel/p4web/pkg-plist
@unexec if cmp -s %D/etc/perforce.conf %D/etc/perforce.conf.default; then rm -f %D/etc/perforce.conf; fi
etc/perforce.conf.default
@exec [ ! -f %B/perforce.conf ] && cp -p %B/%f %B/perforce.conf || true
etc/rc.d/perforce.sh
