# net/poptop/files/openrc-pptpd.in
#!/sbin/openrc-run

name=pptpd
command=%%PREFIX%%/sbin/${name}
pidfile=/var/run/${name}.pid

depend() {
	provide pptpd
}
