# net/kea/files/kea.in
#!/bin/sh
# $FreeBSD$
#
# PROVIDE: kea
# REQUIRE: NETWORK netif routing %%REQ_MYSQL%% %%REQ_PGSQL%%
# KEYWORD: shutdown

#
# Add the following to /etc/rc.conf[.local] to enable this service
#
# kea_enable="YES"
#

. /etc/rc.subr

name=kea
desc="Kea DHCP Server"
rcvar=kea_enable

load_rc_config $name

kea_enable=${kea_enable:-"NO"}

command="%%PREFIX%%/sbin/keactrl"
required_files="%%PREFIX%%/etc/${name}/keactrl.conf"
extra_commands=reload
start_precmd=command_args=start
stop_cmd="$command stop"
status_cmd="$command status"
reload_cmd="$command reload"

run_rc_command "$1"
