# sysutils/bacula-server/pkg-plist
%%SERVER%%share/bacula/query.sql
etc/rc.d/bacula.sh.sample
%%SERVER%%sbin/chio-bacula
%%SERVER%%sbin/bacula-dir
sbin/bacula-fd
%%SERVER%%sbin/bacula-sd
%%SERVER%%sbin/bextract
%%SERVER%%sbin/bls
%%SERVER%%sbin/bscan
%%SERVER%%sbin/btape
sbin/bconsole
%%WXCONS%%sbin/wx-console
%%GNOMECONS%%sbin/gnome-console
%%SERVER%%sbin/dbcheck
%%SERVER%%sbin/bsmtp
%%SERVER%%share/bacula/create_bacula_database
%%SERVER%%share/bacula/create_%%DBTYPE%%_database
%%SERVER%%share/bacula/drop_bacula_database
%%SERVER%%share/bacula/drop_%%DBTYPE%%_database
%%SERVER%%share/bacula/drop_bacula_tables
%%SERVER%%share/bacula/drop_%%DBTYPE%%_tables
%%SERVER%%share/bacula/grant_bacula_privileges
%%SERVER%%share/bacula/grant_%%DBTYPE%%_privileges
%%SERVER%%share/bacula/make_bacula_tables
%%SERVER%%share/bacula/make_%%DBTYPE%%_tables
%%SERVER%%share/bacula/update_bacula_tables
%%SERVER%%share/bacula/update_%%DBTYPE%%_tables
%%SERVER%%share/bacula/delete_catalog_backup
%%SERVER%%share/bacula/make_catalog_backup
%%SERVER%%@dirrm share/bacula
@exec mkdir -p %%BACULA_DIR%%
