# sysutils/bacula-server/pkg-plist
%%DATADIR%%/query.sql
etc/rc.d/z-bacula.sh.sample
sbin/chio-bacula
sbin/bacula-dir
sbin/bacula-sd
sbin/bcopy
sbin/bextract
sbin/bls
sbin/bscan
sbin/btape
%%WXCONS%%sbin/wx-console
%%GNOMECONS%%sbin/gnome-console
sbin/dbcheck
sbin/bsmtp
@unexec if cmp -s %D/share/bacula/mtx-changer.old %D/share/bacula/mtx-changer; then rm -f %D/share/bacula/mtx-changer.old; fi
%%DATADIR%%/mtx-changer
%%DATADIR%%/create_bacula_database
%%DATADIR%%/create_%%DBTYPE%%_database
%%DATADIR%%/drop_bacula_database
%%DATADIR%%/drop_%%DBTYPE%%_database
%%DATADIR%%/drop_bacula_tables
%%DATADIR%%/drop_%%DBTYPE%%_tables
%%DATADIR%%/grant_bacula_privileges
%%DATADIR%%/grant_%%DBTYPE%%_privileges
%%DATADIR%%/make_bacula_tables
%%DATADIR%%/make_%%DBTYPE%%_tables
%%DATADIR%%/update_bacula_tables
%%DATADIR%%/update_%%DBTYPE%%_tables
%%DATADIR%%/delete_catalog_backup
%%DATADIR%%/make_catalog_backup
@dirrm %%DATADIR%%
@unexec /usr/bin/killall badula-sd > /dev/null 2>&1 || true
@unexec /usr/bin/killall bacula-dir > /dev/null 2>&1 || true
@exec mkdir -p %%BACULA_DIR%%
