# sysutils/bacula-server/pkg-plist
%%SERVER%%etc/query.sql
%%SERVER%%etc/rc.d/bacula.sh.sample
%%SERVER%%sbin/chio-bacula
%%SERVER%%sbin/bacula-dir
sbin/bacula-fd
%%SERVER%%sbin/bacula-sd
%%SERVER%%sbin/bextract
%%SERVER%%sbin/bls
%%SERVER%%sbin/bscan
%%SERVER%%sbin/btape
%%SERVER%%sbin/console
%%SERVER%%sbin/dbcheck
%%SERVER%%sbin/smtp
%%SERVER%%share/bacula/create_%%DBTYPE%%_database
%%SERVER%%share/bacula/delete_catalog_backup
%%SERVER%%share/bacula/drop_bacula_tables
%%SERVER%%share/bacula/drop_%%DBTYPE%%_tables
%%SERVER%%share/bacula/make_bacula_tables
%%SERVER%%share/bacula/make_%%DBTYPE%%_tables
%%SERVER%%share/bacula/make_catalog_backup
%%PORTDOCS%%%%DOCSDIR%%/bacula.pdf
%%PORTDOCS%%%%DOCSDIR%%/html-manual/.cvsignore
%%PORTDOCS%%%%DOCSDIR%%/html-manual/Bacula.apf
%%PORTDOCS%%%%DOCSDIR%%/html-manual/add-nav.list
%%PORTDOCS%%%%DOCSDIR%%/html-manual/add-nav.pl
%%PORTDOCS%%%%DOCSDIR%%/html-manual/autochangers.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/autochangers.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/back.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/bacula-features.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/bacula-features.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/bacula-logo.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/bacula-styles.css
%%PORTDOCS%%%%DOCSDIR%%/html-manual/bacula.book
%%PORTDOCS%%%%DOCSDIR%%/html-manual/bacula.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/bacula.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/bacula_template.inc
%%PORTDOCS%%%%DOCSDIR%%/html-manual/bclogo.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/bootstrap.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/bootstrap.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/bugs.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/bugs.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/catalog.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/catalog.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/catmaintenance.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/catmaintenance.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/configure.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/configure.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/console.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/console.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/consoleconf.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/consoleconf.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/daemonprotocol.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/daemonprotocol.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/developers.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/developers.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/dirdconf.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/dirdconf.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/director.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/director.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/disk.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/disk.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/error.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/faq.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/faq.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/file.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/file.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/file.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/filedconf.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/filedconf.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/flow.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/flow.jpeg
%%PORTDOCS%%%%DOCSDIR%%/html-manual/gpl.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/gpl.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/gui-interface.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/gui-interface.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/home.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/idle.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/imgdot-1x1-000000.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/imgdot-1x1-transp-ffffff.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/index.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/index.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/install.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/install.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/internaldb.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/internaldb.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/internals.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/internals.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/kaboom.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/kaboom.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/lesser.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/lesser.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/license.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/license.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/md5.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/md5.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/md5.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/mediaformat.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/mediaformat.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/mempool.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/mempool.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/messagesres.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/messagesres.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/mysql.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/mysql.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/netprotocol.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/netprotocol.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/next.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/philosophical-gnu-sm.jpg
%%PORTDOCS%%%%DOCSDIR%%/html-manual/porting.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/porting.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/progs.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/progs.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/projects.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/projects.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/quickstart.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/quickstart.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/recycling.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/recycling.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/regression.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/regression.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/rescue.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/rescue.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/restore.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/restore.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/running.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/running.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/running.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/security.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/security.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/smartall.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/smartall.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/smartall.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/sqlite.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/sqlite.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/state.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/state.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/status.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/storage.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/storage.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/storedconf.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/storedconf.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/strategies.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/strategies.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/tapetesting.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/tapetesting.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/thanks.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/thanks.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/thanks.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/tips.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/tips.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/vars.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/vars.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/verify.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/verify.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/waiting.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/win32.html
%%PORTDOCS%%%%DOCSDIR%%/html-manual/win32.wml
%%PORTDOCS%%%%DOCSDIR%%/html-manual/wininstall-1.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/wininstall-2.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/wininstall-3.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/wininstall-4.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/wininstall-5.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/wininstall-6.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/wininstall-7.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/wininstall-8.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/wininstall-9.gif
%%PORTDOCS%%%%DOCSDIR%%/html-manual/working.gif
%%PORTDOCS%%@dirrm %%DOCSDIR%%/html-manual
%%PORTDOCS%%@dirrm %%DOCSDIR%%
%%SERVER%%@dirrm share/bacula
