# databases/postgresql73-server/pkg-plist-server
bin/postmaster
bin/postgres
etc/periodic/daily/502.pgsql
etc/rc.d/010.pgsql.sh
lib/postgresql/ascii_and_mic.so
lib/postgresql/cyrillic_and_mic.so
lib/postgresql/euc_cn_and_mic.so
lib/postgresql/euc_jp_and_sjis.so
lib/postgresql/euc_kr_and_mic.so
lib/postgresql/euc_tw_and_big5.so
lib/postgresql/latin2_and_win1250.so
lib/postgresql/latin_and_mic.so
lib/postgresql/utf8_and_ascii.so
lib/postgresql/utf8_and_big5.so
lib/postgresql/utf8_and_cyrillic.so
lib/postgresql/utf8_and_euc_cn.so
lib/postgresql/utf8_and_euc_jp.so
lib/postgresql/utf8_and_euc_kr.so
lib/postgresql/utf8_and_euc_tw.so
lib/postgresql/utf8_and_gb18030.so
lib/postgresql/utf8_and_gbk.so
lib/postgresql/utf8_and_iso8859.so
lib/postgresql/utf8_and_iso8859_1.so
lib/postgresql/utf8_and_johab.so
lib/postgresql/utf8_and_sjis.so
lib/postgresql/utf8_and_tcvn.so
lib/postgresql/utf8_and_uhc.so
lib/postgresql/utf8_and_win1250.so
lib/postgresql/utf8_and_win1256.so
lib/postgresql/utf8_and_win874.so
lib/postgresql/plpgsql.so
%%DOCSDIR%%/README-server
@unexec rmdir %D/%%DOCSDIR%% 2>/dev/null || true
share/postgresql/postgres.bki
share/postgresql/postgres.description
share/postgresql/pg_hba.conf.sample
share/postgresql/pg_ident.conf.sample
share/postgresql/postgresql.conf.sample
share/postgresql/conversion_create.sql
%%GETTEXT%%share/locale/cs/LC_MESSAGES/postgres.mo
%%GETTEXT%%share/locale/de/LC_MESSAGES/postgres.mo
%%GETTEXT%%share/locale/es/LC_MESSAGES/postgres.mo
%%GETTEXT%%share/locale/hr/LC_MESSAGES/postgres.mo
%%GETTEXT%%share/locale/hu/LC_MESSAGES/postgres.mo
%%GETTEXT%%share/locale/ru/LC_MESSAGES/postgres.mo
%%GETTEXT%%share/locale/sv/LC_MESSAGES/postgres.mo
%%GETTEXT%%share/locale/tr/LC_MESSAGES/postgres.mo
%%GETTEXT%%share/locale/zh_CN/LC_MESSAGES/postgres.mo
%%GETTEXT%%share/locale/zh_TW/LC_MESSAGES/postgres.mo
%%GETTEXT%%@unexec rmdir %D/share/locale/cs/LC_MESSAGES 2>/dev/null || true
%%GETTEXT%%@unexec rmdir %D/share/locale/cs 2>/dev/null || true
%%GETTEXT%%@unexec rmdir %D/share/locale/zh_CN/LC_MESSAGES 2>/dev/null || true
%%GETTEXT%%@unexec rmdir %D/share/locale/zh_CN 2>/dev/null || true
%%GETTEXT%%@unexec rmdir %D/share/locale/zh_TW/LC_MESSAGES 2>/dev/null || true
%%GETTEXT%%@unexec rmdir %D/share/locale/zh_TW 2>/dev/null || true
@unexec if cmp -s ~pgsql/.profile %D/share/postgresql/dot.profile.dist; then rm -f ~pgsql/.profile; fi
@unexec if cmp -s ~pgsql/.cshrc %D/share/postgresql/dot.cshrc.dist; then rm -f ~pgsql/.cshrc; fi
@unexec [ -f ~pgsql/dot.profile.dist ] && rm ~pgsql/dot.profile.dist
@unexec [ -f ~pgsql/dot.cshrc.dist ] && rm ~pgsql/dot.cshrc.dist
share/postgresql/dot.profile.dist
share/postgresql/dot.cshrc.dist
@exec [ -f %D/share/postgresql/dot.profile.dist ] && cp %D/share/postgresql/dot.profile.dist ~pgsql/.
@exec [ -f %D/share/postgresql/dot.cshrc.dist ] && cp %D/share/postgresql/dot.cshrc.dist ~pgsql/.
@exec [ -f ~pgsql/.profile ] || (cp %D/share/postgresql/dot.profile.dist ~pgsql/.profile && chown pgsql:pgsql ~pgsql/.profile)
@exec [ -f ~pgsql/.cshrc ] || (cp %D/share/postgresql/dot.cshrc.dist ~pgsql/.cshrc && chown pgsql:pgsql ~pgsql/.cshrc)
@unexec echo "Warning: If you will *NOT* use this package anymore, please remove the pgsql user manually."
@unexec rmdir %D/share/postgresql 2>/dev/null || true
@unexec rmdir %D/etc/periodic/daily 2>/dev/null || true
@unexec rmdir %D/etc/periodic 2>/dev/null || true
@unexec rmdir %D/lib/postgresql 2>/dev/null || true
@unexec rmdir %D/pgsql 2>/dev/null || true
