# databases/msql3/pkg-plist
bin/cpi_stat
bin/msql
bin/msql3_broker
bin/msqladmin
bin/msqldump
bin/msqlexplain
bin/msqlexport
bin/msqlimport
bin/relshow
bin/run_daemon
sbin/msql3d
etc/msql3/msql.conf
etc/rc.d/msql3.sh
include/msql.h
include/regexp.h
lib/libmsql.a
lib/libregexp.a
%%DATADIR%%/makegen/directory.mm
%%DATADIR%%/makegen/dll_os2.mm
%%DATADIR%%/makegen/install.mm
%%DATADIR%%/makegen/lex.mm
%%DATADIR%%/makegen/libinstall.mm
%%DATADIR%%/makegen/library_os2.mm
%%DATADIR%%/makegen/library_unix.mm
%%DATADIR%%/makegen/makegen
%%DATADIR%%/makegen/makegen.cf
%%DATADIR%%/makegen/object.mm
%%DATADIR%%/makegen/program.mm
%%DATADIR%%/makegen/site.mm
%%DATADIR%%/makegen/touch.mm
%%DATADIR%%/makegen/yacc.mm
%%PORTDOCS%%%%DOCSDIR%%/BUGS
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/RELEASE_NOTES
%%PORTDOCS%%%%DOCSDIR%%/License
%%PORTDOCS%%%%DOCSDIR%%/Manual
%%PORTDOCS%%%%DOCSDIR%%/Perl
%%PORTDOCS%%%%DOCSDIR%%/acl.txt
@dirrm etc/msql3
@dirrm include/common
@dirrm %%DATADIR%%/makegen
@dirrm %%DATADIR%%
%%PORTDOCS%%@dirrm %%DOCSDIR%%
