# databases/msql3/pkg-plist
bin/cpi_stat
bin/msql
bin/msql3_broker
bin/msqladmin
bin/msqldump
bin/msqlexplain
bin/msqlexport
bin/msqlimport
bin/relshow
bin/run_daemon
sbin/msql3d
etc/msql3/msql.conf
etc/rc.d/msql3.sh
include/msql.h
include/regexp.h
lib/libmsql.a
lib/libregexp.a
share/msql3/makegen/directory.mm
share/msql3/makegen/dll_os2.mm
share/msql3/makegen/install.mm
share/msql3/makegen/lex.mm
share/msql3/makegen/libinstall.mm
share/msql3/makegen/library_os2.mm
share/msql3/makegen/library_unix.mm
share/msql3/makegen/makegen
share/msql3/makegen/makegen.cf
share/msql3/makegen/object.mm
share/msql3/makegen/program.mm
share/msql3/makegen/site.mm
share/msql3/makegen/touch.mm
share/msql3/makegen/yacc.mm
share/doc/msql3/BUGS
share/doc/msql3/README
share/doc/msql3/RELEASE_NOTES
share/doc/msql3/License
share/doc/msql3/Manual
share/doc/msql3/Perl
share/doc/msql3/acl.txt
@dirrm etc/msql3
@dirrm include/common
@dirrm share/msql3/makegen
@dirrm share/msql3
@dirrm share/doc/msql3
