# databases/msql3/pkg-plist
bin/lite
bin/msql
bin/msqladmin
bin/msqldump
bin/msqlexport
bin/msqlimport
bin/relshow
bin/run_daemon
bin/w3-auth
bin/w3-msql
etc/msql2/msql.acl.sample
etc/msql2/msql.conf
etc/rc.d/msql2.sh
include/common/config.h
include/common/portability.h
include/lite.h
include/mod_msql.h
include/mod_std.h
include/msql.h
lib/liblite.a
lib/liblite.so.1
@exec ln -sf %f %B/liblite.so
@unexec rm -f %B/liblite.so
lib/libmsql.a
lib/libmsql.so.1
@exec ln -sf %f %B/libmsql.so
@unexec rm -f %B/libmsql.so
lib/msql/lite_lex.o
lib/msql/mod_lite.o
lib/msql/mod_msql.o
lib/msql/mod_std.o
libexec/msql/mod_lite.so
libexec/msql/mod_msql.so
libexec/msql/mod_std.so
sbin/msql2d
share/doc/msql/HISTORY
share/doc/msql/INSTALL
share/doc/msql/License
share/doc/msql/MSQL_BOOK
share/doc/msql/NEW_FEATURES
share/doc/msql/README
share/doc/msql/README.sco
share/doc/msql/RELEASE_NOTES
share/doc/msql/invoice.ps.gz
share/doc/msql/invoice.txt
share/doc/msql/manual-html/Image1.gif
share/doc/msql/manual-html/Image2.gif
share/doc/msql/manual-html/Image3.gif
share/doc/msql/manual-html/Image4.gif
share/doc/msql/manual-html/Image5.gif
share/doc/msql/manual-html/Image6.gif
share/doc/msql/manual-html/deer.gif
share/doc/msql/manual-html/manual.html
share/doc/msql/manual.ps.gz
share/doc/msql/post-install-notes
share/msql/makegen/directory.mm
share/msql/makegen/dll_os2.mm
share/msql/makegen/install.mm
share/msql/makegen/lex.mm
share/msql/makegen/libinstall.mm
share/msql/makegen/library_os2.mm
share/msql/makegen/library_unix.mm
share/msql/makegen/makegen
share/msql/makegen/makegen.cf
share/msql/makegen/makegen.cf.in
share/msql/makegen/object.mm
share/msql/makegen/program.mm
share/msql/makegen/site.mm
share/msql/makegen/slibinstall.mm
share/msql/makegen/slibrary_unix.mm
share/msql/makegen/touch.mm
share/msql/makegen/yacc.mm
share/msql/misc/setup_w3auth
share/msql/www/bookmarks/README
share/msql/www/bookmarks/Welcome.html
share/msql/www/bookmarks/add.html
share/msql/www/bookmarks/back.gif
share/msql/www/bookmarks/book.gif
share/msql/www/bookmarks/bookmarks.html
share/msql/www/bookmarks/create.gif
share/msql/www/bookmarks/create.html
share/msql/www/bookmarks/delete.gif
share/msql/www/bookmarks/delete.html
share/msql/www/bookmarks/dump.msql
share/msql/www/bookmarks/folder.gif
share/msql/www/bookmarks/home.gif
share/msql/www/bookmarks/mark-big.gif
share/msql/www/bookmarks/remove.html
share/msql/www/bookmarks/search.gif
share/msql/www/bookmarks/search.html
share/msql/www/bookmarks/search_result.html
share/msql/www/bookmarks/setup_bookmark
share/msql/www/bookmarks/wood.gif
share/msql/www/graphics/add.gif
share/msql/www/graphics/area.gif
share/msql/www/graphics/banner.gif
share/msql/www/graphics/blank.gif
share/msql/www/graphics/delete.gif
share/msql/www/graphics/edit.gif
share/msql/www/graphics/group.gif
share/msql/www/graphics/hughes.gif
share/msql/www/graphics/logout.gif
share/msql/www/graphics/priv.gif
share/msql/www/graphics/user.gif
share/msql/www/graphics/view.gif
@dirrm etc/msql2
@dirrm include/common
@dirrm lib/msql
@dirrm libexec/msql
@dirrm share/doc/msql/manual-html
@dirrm share/doc/msql
@dirrm share/msql/makegen
@dirrm share/msql/misc
@dirrm share/msql/www/bookmarks
@dirrm share/msql/www/graphics
@dirrm share/msql/www
@dirrm share/msql
