# sysutils/ganglia-monitor-core/pkg-plist
bin/gstat
bin/gmetric
etc/gmond.conf.sample
%%GMETAD%%etc/gmetad.conf.sample
etc/rc.d/gmond.sh.sample
%%GMETAD%%etc/rc.d/gmetad.sh.sample
include/ganglia.h
include/ascii.h
include/asciitab.h
include/daemon_inetd.h
include/daemon_init.h
include/debug_msg.h
include/error.h
include/file.h
include/hash.h
include/iasciitab.h
include/latin1tab.h
include/llist.h
include/my_inet_ntop.h
include/nametab.h
include/net.h
include/rdwr.h
include/utf8tab.h
include/xmldef.h
include/xmlparse.h
include/xmlrole.h
include/xmltok.h
include/xmltok_impl.h
include/gexec_funcs.h
include/fsusage.h
include/gexec_process.h
include/gmond_config.h
lib/libganglia.so.2
lib/libganglia.so
lib/libganglia.la
lib/libganglia.a
%%GMETAD%%sbin/gmetad
sbin/gmond
@dirrm include/ganglia
@unexec rmdir /var/db/ganglia/rrds /var/db/ganglia 2>/dev/null || true
