# security/wpa_supplicant/files/openrc-wpa_supplicant.in
#!/sbin/openrc-run
# Copyright (c) 2007-2015 The OpenRC Authors.
# See the Authors file at the top-level directory of this distribution and
# https://github.com/OpenRC/openrc/blob/master/AUTHORS
#
# This file is part of OpenRC. It is subject to the license terms in
# the LICENSE file found in the top-level directory of this
# distribution and at https://github.com/OpenRC/openrc/blob/master/LICENSE
# This file may not be copied, modified, propagated, or distributed
# except according to the terms contained in the LICENSE file.

BIN=%%PREFIX%%/sbin/wpa_supplicant
: ${wpa_supplicant_conf:=/etc/wpa_supplicant.conf}

wpa=${RC_SVCNAME##*.}
if [ -n "$wpa" -a "$wpa" != "wpa_supplicant" ]; then
        wpa_device="$wpa"
        pidfile="/var/run/wpa_supplicant-${wpa_device}.pid"
fi

wpa_args="$wpa_supplicant_args -B -c$wpa_supplicant_conf"
name="wpa_supplicant"
[ -n "$wpa_device" ] && name="$name ($wpa_device)"

depend()
{
	keyword -shutdown -stop -jail
}

find_wireless()
{
	local _devs
	for i in `ifconfig -l 2>/dev/null`
	do
		echo $i | grep -q "^wlan"
		if [ $? -ne 0 ] ; then continue; fi
		if wpaif $i; then
			_devs="$_devs $i"
		fi
	done

	echo "$_devs"
}

start_wpa_service()
{
	device="$1"
        if [ ! -e "/etc/wpa_supplicant.conf" ] ; then
		echo "ctrl_interface=/var/run/wpa_supplicant" > /etc/wpa_supplicant.conf
		chmod 600 /etc/wpa_supplicant.conf
        fi

	wpa_args="$wpa_args -i${device} -P ${pidfile}"

	start-stop-daemon --start --exec $BIN -p ${pidfile} -- ${wpa_args}
}

start()
{
	if [ -n "$wpa_device" ] ; then
		if [ -z "$OPENRC_SILENT" ] ; then
			ebegin "Starting $name"
		fi
		start_wpa_service "$wpa_device"
		stat=$?
		if [ -z "$OPENRC_SILENT" ] ; then
			eend $stat
		fi
		return $stat
	fi

	# Start any sub-devices
	for i in $(find_wireless)
	do
		export OPENRC_SILENT=TRUE
		ln -fs %%PREFIX%%/etc/init.d/wpa_supplicant %%PREFIX%%/etc/init.d/wpa_supplicant.${i}
		/sbin/service wpa_supplicant.${i} start
	done
	return 0
}

stop_wpa_service()
{
	# Stop DHCPCD for the wireless device
	service dhcpd.${1} stop >/dev/null 2>/dev/null
	start-stop-daemon --stop --exec $BIN -p ${pidfile} -- ${wpa_args}
}

stop()
{
	if [ -n "$wpa_device" ] ; then
		if [ -z "$OPENRC_SILENT" ] ; then
			ebegin "Stopping $name"
		fi
		stop_wpa_service "$wpa_device"
		stat=$?
		if [ -z "$OPENRC_SILENT" ] ; then
			eend $stat
		fi
		return $stat
	fi

	# Stopping any sub-devices
	for i in $(find_wireless)
	do
		export OPENRC_SILENT=TRUE
		ln -fs %%PREFIX%%/etc/init.d/wpa_supplicant %%PREFIX%%/etc/init.d/wpa_supplicant.${i}
		/sbin/service wpa_supplicant.${i} stop
	done
	return 0
}
