# net-p2p/datum_gateway/pkg-plist
bin/datum_gateway
%%ETCDIR%%_config.json.sample
etc/rc.d/datum_gateway
share/doc/DATUM/README.md
share/doc/DATUM/doc/DATUM_recommended_setup-network_diagram.svg
share/doc/DATUM/example_datum_gateway_config.json
share/doc/DATUM/usernames.md
@dir(nobody,nogroup) /var/run/datum_gateway
