# mail/dbmail/files/openrc-dbmail-imapd.in
#!/sbin/openrc-run

name=dbmail_imapd
command=%%PREFIX%%/sbin/dbmail-imapd
pidfile=/var/run/dbmail/dbmail-imapd.pid
required_files=%%PREFIX%%/etc/dbmail.conf

depend() {
	keyword -shutdown
}
