# net-mgmt/nagios12/pkg-plist
@exec mkdir -p %D/var/netsaint/archives
@exec mkdir -p %D/var/netsaint/rw
@exec chmod 775 %D/var/netsaint %D/var/netsaint/archives %D/var/netsaint/rw
@exec chown netsaint:netsaint %D/var/netsaint %D/var/netsaint/archives
@exec chown netsaint:nogroup %D/var/netsaint/rw
bin/netsaint
etc/netsaint/commands.cfg.dist
etc/netsaint/hosts.cfg.dist
etc/netsaint/netsaint.cfg.dist
etc/netsaint/nscgi.cfg.dist
etc/netsaint/resource.cfg.dist
etc/rc.d/netsaint.sh
share/netsaint/docs/about.html
share/netsaint/docs/addons.html
share/netsaint/docs/cgiauth.html
share/netsaint/docs/cgis.html
share/netsaint/docs/checkscheduling.html
share/netsaint/docs/clusters.html
share/netsaint/docs/commandfile.html
share/netsaint/docs/config.html
share/netsaint/docs/configcgi.html
share/netsaint/docs/confighost.html
share/netsaint/docs/configmain.html
share/netsaint/docs/dependencies.html
share/netsaint/docs/developer/commentfile.html
share/netsaint/docs/developer/images/xcdapi.gif
share/netsaint/docs/developer/images/xedapi.gif
share/netsaint/docs/developer/images/xodapi.gif
share/netsaint/docs/developer/images/xrdapi.gif
share/netsaint/docs/developer/images/xsdapi.gif
share/netsaint/docs/developer/index.html
share/netsaint/docs/developer/pluginhowto.html
share/netsaint/docs/developer/stateretentionfile.html
share/netsaint/docs/developer/statusfile.html
share/netsaint/docs/developer/xcd.html
share/netsaint/docs/developer/xdata.html
share/netsaint/docs/developer/xed.html
share/netsaint/docs/developer/xod.html
share/netsaint/docs/developer/xrd.html
share/netsaint/docs/developer/xsd.html
share/netsaint/docs/distributed.html
share/netsaint/docs/downtime.html
share/netsaint/docs/embeddedperl.html
share/netsaint/docs/escalations.html
share/netsaint/docs/eventhandlers.html
share/netsaint/docs/extcommands.html
share/netsaint/docs/faqs.html
share/netsaint/docs/flapping.html
share/netsaint/docs/funstuff.html
share/netsaint/docs/images/activepassive.gif
share/netsaint/docs/images/avail_cgi.gif
share/netsaint/docs/images/checktiming.gif
share/netsaint/docs/images/cmd_cgi.gif
share/netsaint/docs/images/config_cgi.gif
share/netsaint/docs/images/distributed.gif
share/netsaint/docs/images/extinfo_cgi_a.gif
share/netsaint/docs/images/extinfo_cgi_b.gif
share/netsaint/docs/images/extinfo_cgi_c.gif
share/netsaint/docs/images/extinfo_cgi_d.gif
share/netsaint/docs/images/extinfo_cgi_e.gif
share/netsaint/docs/images/history_cgi.gif
share/netsaint/docs/images/indirecthostcheck.gif
share/netsaint/docs/images/indirectsvccheck.gif
share/netsaint/docs/images/indirectsvccheck2.gif
share/netsaint/docs/images/interleaved1.gif
share/netsaint/docs/images/interleaved2.gif
share/netsaint/docs/images/network-heirarchy.gif
share/netsaint/docs/images/network-outage1.gif
share/netsaint/docs/images/network-outage2.gif
share/netsaint/docs/images/noninterleaved1.gif
share/netsaint/docs/images/noninterleaved2.gif
share/netsaint/docs/images/notification-process.gif
share/netsaint/docs/images/notifications_cgi.gif
share/netsaint/docs/images/outages_cgi.gif
share/netsaint/docs/images/physical-network.gif
share/netsaint/docs/images/plugintheory.gif
share/netsaint/docs/images/redundancy-1.gif
share/netsaint/docs/images/redundancy-2.gif
share/netsaint/docs/images/service-dependencies.gif
share/netsaint/docs/images/showlog_cgi.gif
share/netsaint/docs/images/statetransitions.gif
share/netsaint/docs/images/statetransitions2.gif
share/netsaint/docs/images/status_cgi_a.gif
share/netsaint/docs/images/status_cgi_b.gif
share/netsaint/docs/images/statusmap_cgi.gif
share/netsaint/docs/images/statuswml_cgi.gif
share/netsaint/docs/images/statuswrl_cgi.gif
share/netsaint/docs/images/tac_cgi.gif
share/netsaint/docs/images/trends_cgi.gif
share/netsaint/docs/index.html
share/netsaint/docs/indirectchecks.html
share/netsaint/docs/installing.html
share/netsaint/docs/installweb.html
share/netsaint/docs/macros.html
share/netsaint/docs/makehtml.html
share/netsaint/docs/networkoutages.html
share/netsaint/docs/networkreachability.html
share/netsaint/docs/notifications.html
share/netsaint/docs/parallelization.html
share/netsaint/docs/passivechecks.html
share/netsaint/docs/perfdata.html
share/netsaint/docs/pluginhowto.html
share/netsaint/docs/plugins.html
share/netsaint/docs/plugintheory.html
share/netsaint/docs/programmodes.html
share/netsaint/docs/redundancy.html
share/netsaint/docs/robots.txt
share/netsaint/docs/security.html
share/netsaint/docs/starting.html
share/netsaint/docs/statetypes.html
share/netsaint/docs/statuslevels.html
share/netsaint/docs/stoprestart.html
share/netsaint/docs/theory.html
share/netsaint/docs/timeperiods.html
share/netsaint/docs/tuning.html
share/netsaint/docs/verifyconfig.html
share/netsaint/docs/volatileservices.html
share/netsaint/docs/whatsnew.html
share/netsaint/docs/xdata-db.html
share/netsaint/images/ack.gif
share/netsaint/images/active.gif
share/netsaint/images/bsd_logofullsize.gif
share/netsaint/images/bsd_sblogo.gif
share/netsaint/images/command.gif
share/netsaint/images/comment.gif
share/netsaint/images/config.gif
share/netsaint/images/contact.gif
share/netsaint/images/contactgroup.gif
share/netsaint/images/critical.gif
share/netsaint/images/delay.gif
share/netsaint/images/delete.gif
share/netsaint/images/disabled.gif
share/netsaint/images/down.gif
share/netsaint/images/downtime.gif
share/netsaint/images/enabled.gif
share/netsaint/images/extinfo.gif
share/netsaint/images/flapping.gif
share/netsaint/images/greendot.gif
share/netsaint/images/history.gif
share/netsaint/images/host.gif
share/netsaint/images/hostevent.gif
share/netsaint/images/info.gif
share/netsaint/images/left.gif
share/netsaint/images/logofullsize.gif
share/netsaint/images/logos/aix.gd2
share/netsaint/images/logos/aix.gif
share/netsaint/images/logos/amiga.gd2
share/netsaint/images/logos/amiga.gif
share/netsaint/images/logos/apple.gd2
share/netsaint/images/logos/apple.gif
share/netsaint/images/logos/beos.gd2
share/netsaint/images/logos/beos.gif
share/netsaint/images/logos/caldera.gd2
share/netsaint/images/logos/caldera.gif
share/netsaint/images/logos/cat1900.gd2
share/netsaint/images/logos/cat2900.gd2
share/netsaint/images/logos/cat5000.gd2
share/netsaint/images/logos/debian.gd2
share/netsaint/images/logos/debian.gif
share/netsaint/images/logos/freebsd40.gd2
share/netsaint/images/logos/freebsd40.gif
share/netsaint/images/logos/hp-printer40.gd2
share/netsaint/images/logos/hp-printer40.gif
share/netsaint/images/logos/hpux.gd2
share/netsaint/images/logos/hpux.gif
share/netsaint/images/logos/irix.gd2
share/netsaint/images/logos/irix.gif
share/netsaint/images/logos/linux40.gd2
share/netsaint/images/logos/linux40.gif
share/netsaint/images/logos/logo.gd2
share/netsaint/images/logos/mac40.gd2
share/netsaint/images/logos/mac40.gif
share/netsaint/images/logos/mandrake.gd2
share/netsaint/images/logos/mandrake.gif
share/netsaint/images/logos/next.gd2
share/netsaint/images/logos/next.gif
share/netsaint/images/logos/netsaint.gd2
share/netsaint/images/logos/ng-switch40.gd2
share/netsaint/images/logos/ng-switch40.gif
share/netsaint/images/logos/novell40.gd2
share/netsaint/images/logos/novell40.gif
share/netsaint/images/logos/openbsd.gd2
share/netsaint/images/logos/openbsd.gif
share/netsaint/images/logos/redhat.gd2
share/netsaint/images/logos/redhat.gif
share/netsaint/images/logos/router40.gd2
share/netsaint/images/logos/router40.gif
share/netsaint/images/logos/slackware.gd2
share/netsaint/images/logos/slackware.gif
share/netsaint/images/logos/stampede.gd2
share/netsaint/images/logos/stampede.gif
share/netsaint/images/logos/storm.gd2
share/netsaint/images/logos/storm.gif
share/netsaint/images/logos/sun40.gd2
share/netsaint/images/logos/sun40.gif
share/netsaint/images/logos/sunlogo.gd2
share/netsaint/images/logos/sunlogo.gif
share/netsaint/images/logos/switch40.gd2
share/netsaint/images/logos/switch40.gif
share/netsaint/images/logos/turbolinux.gd2
share/netsaint/images/logos/turbolinux.gif
share/netsaint/images/logos/ultrapenguin.gd2
share/netsaint/images/logos/ultrapenguin.gif
share/netsaint/images/logos/unicos.gd2
share/netsaint/images/logos/unicos.gif
share/netsaint/images/logos/unknown.gd2
share/netsaint/images/logos/win40.gd2
share/netsaint/images/logos/win40.gif
share/netsaint/images/logos/yellowdog.gd2
share/netsaint/images/logos/yellowdog.gif
share/netsaint/images/logrotate.gif
share/netsaint/images/logs.gif
share/netsaint/images/ndisabled.gif
share/netsaint/images/noack.gif
share/netsaint/images/notes.gif
share/netsaint/images/notify.gif
share/netsaint/images/orangedot.gif
share/netsaint/images/parentup.gif
share/netsaint/images/recovery.gif
share/netsaint/images/restart.gif
share/netsaint/images/right.gif
share/netsaint/images/sbbg.gif
share/netsaint/images/sbconfig.gif
share/netsaint/images/sbgeneral.gif
share/netsaint/images/sblogo.gif
share/netsaint/images/sbmonitor.gif
share/netsaint/images/serviceevent.gif
share/netsaint/images/sideback.jpg
share/netsaint/images/standby.gif
share/netsaint/images/start.gif
share/netsaint/images/status.gif
share/netsaint/images/status2.gif
share/netsaint/images/status3.gif
share/netsaint/images/status4.gif
share/netsaint/images/stop.gif
share/netsaint/images/tacdisabled.jpg
share/netsaint/images/tacenabled.jpg
share/netsaint/images/thermcrit.gif
share/netsaint/images/thermok.gif
share/netsaint/images/thermwarn.gif
share/netsaint/images/trends.gif
share/netsaint/images/trendshost.png
share/netsaint/images/trendssvc.png
share/netsaint/images/unknown.gif
share/netsaint/images/up.gif
share/netsaint/images/warning.gif
share/netsaint/images/zoom1.gif
share/netsaint/images/zoom2.gif
share/netsaint/index.html
share/netsaint/main.html
share/netsaint/media/critical.wav
share/netsaint/media/hostdown.wav
share/netsaint/media/warning.wav
share/netsaint/robots.txt
share/netsaint/sbin/avail.cgi
share/netsaint/sbin/cmd.cgi
share/netsaint/sbin/config.cgi
share/netsaint/sbin/extinfo.cgi
share/netsaint/sbin/history.cgi
share/netsaint/sbin/ministatus.cgi
share/netsaint/sbin/notifications.cgi
share/netsaint/sbin/outages.cgi
share/netsaint/sbin/showlog.cgi
share/netsaint/sbin/status.cgi
share/netsaint/sbin/statusmap.cgi
share/netsaint/sbin/statuswml.cgi
share/netsaint/sbin/statuswrl.cgi
share/netsaint/sbin/tac.cgi
share/netsaint/sbin/trends.cgi
share/netsaint/side.html
share/netsaint/stylesheets/avail.css
share/netsaint/stylesheets/cmd.css
share/netsaint/stylesheets/config.css
share/netsaint/stylesheets/extinfo.css
share/netsaint/stylesheets/history.css
share/netsaint/stylesheets/notifications.css
share/netsaint/stylesheets/outages.css
share/netsaint/stylesheets/showlog.css
share/netsaint/stylesheets/status.css
share/netsaint/stylesheets/statusmap.css
share/netsaint/stylesheets/tac.css
share/netsaint/stylesheets/trends.css
@unexec rmdir %D/etc/netsaint 2>/dev/null || true
@dirrm share/netsaint/docs/developer/images
@dirrm share/netsaint/docs/developer
@dirrm share/netsaint/docs/images
@dirrm share/netsaint/docs
@dirrm share/netsaint/images/logos
@dirrm share/netsaint/images
@dirrm share/netsaint/media
@dirrm share/netsaint/sbin
@dirrm share/netsaint/stylesheets
@dirrm share/netsaint
@unexec rmdir %D/var/netsaint/archives 2>/dev/null || true
@unexec rmdir %D/var/netsaint/rw 2>/dev/null || true
@unexec rmdir %D/var/netsaint 2>/dev/null || true
@unexec rmdir %D/var 2>/dev/null || true
