# sysutils/LPRng/pkg-plist
bin/cancel
bin/lp
bin/lpq
bin/lpr
bin/lprm
bin/lpstat
etc/lpd.conf.sample
etc/lpd.perms.sample
etc/lprng.sh
etc/printcap.sample
etc/rc.d/lprng.sh
libexec/filters/lpbanner
libexec/filters/lpf
libexec/filters/pclbanner
libexec/filters/psbanner
sbin/checkpc
sbin/lpc
sbin/lpd
sbin/lpraccnt
share/doc/LPRng/CHANGES
share/doc/LPRng/LISA98.ppt
share/doc/LPRng/LPRng-HOWTO.dsl
share/doc/LPRng/LPRng-HOWTO.dtd
share/doc/LPRng/LPRng-HOWTO.html
share/doc/LPRng/LPRng-HOWTO.sgml
share/doc/LPRng/LPRng.html
share/doc/LPRng/LPRngLOGO100x30.gif
share/doc/LPRng/LPRngLOGO2.gif
share/doc/LPRng/LPRngLOGO334x100.gif
share/doc/LPRng/LPRng_20logo.gif
share/doc/LPRng/Makefile
share/doc/LPRng/Makefile.in
share/doc/LPRng/REDIRECT.LPRng.html
share/doc/LPRng/abnormalterm.htm
share/doc/LPRng/accountingref.htm
share/doc/LPRng/ah.htm
share/doc/LPRng/aix.htm
share/doc/LPRng/allowusersetting.htm
share/doc/LPRng/append-z.htm
share/doc/LPRng/appletalk.htm
share/doc/LPRng/appsocket.htm
share/doc/LPRng/auth.htm
share/doc/LPRng/authref.htm
share/doc/LPRng/bannerprinting.htm
share/doc/LPRng/bg.gif
share/doc/LPRng/bk.htm
share/doc/LPRng/chainingfilters.htm
share/doc/LPRng/checkfornonprintable.htm
share/doc/LPRng/checkpc.htm
share/doc/LPRng/classpriority.htm
share/doc/LPRng/cm.htm
share/doc/LPRng/collateindex.pl
share/doc/LPRng/configfile.htm
share/doc/LPRng/controlfiles.htm
share/doc/LPRng/controllingprintqueue.htm
share/doc/LPRng/debugging.htm
share/doc/LPRng/defaultrmrp.htm
share/doc/LPRng/destinations.htm
share/doc/LPRng/dict
share/doc/LPRng/exitcodes.htm
share/doc/LPRng/f36.htm
share/doc/LPRng/faq.htm
share/doc/LPRng/faqref.htm
share/doc/LPRng/filterenvironmentvars.htm
share/doc/LPRng/filterexitcodes.htm
share/doc/LPRng/filteroptions.htm
share/doc/LPRng/filters.htm
share/doc/LPRng/forcefqdnhostname.htm
share/doc/LPRng/forcelocalhost.htm
share/doc/LPRng/genindex.sgml
share/doc/LPRng/i11210.htm
share/doc/LPRng/i11252.htm
share/doc/LPRng/ifhp-VERSION
share/doc/LPRng/ifhp.htm
share/doc/LPRng/index.htm
share/doc/LPRng/installation.htm
share/doc/LPRng/introduction.htm
share/doc/LPRng/jaggies.htm
share/doc/LPRng/jobcompletion.htm
share/doc/LPRng/jobcopies.htm
share/doc/LPRng/jobfileformatconversion.htm
share/doc/LPRng/jobfiles.htm
share/doc/LPRng/jobformatsandfilterselection.htm
share/doc/LPRng/jobremoval.htm
share/doc/LPRng/jobsteps.htm
share/doc/LPRng/jobtransfer.htm
share/doc/LPRng/kerberos.htm
share/doc/LPRng/license.enc
share/doc/LPRng/license.htm
share/doc/LPRng/license.txt
share/doc/LPRng/linkbanner.gif
share/doc/LPRng/links.html
share/doc/LPRng/linux.htm
share/doc/LPRng/lisa2000.jpg
share/doc/LPRng/loadbalance.htm
share/doc/LPRng/logfiles.htm
share/doc/LPRng/lpc.htm
share/doc/LPRng/lpcredirect.htm
share/doc/LPRng/lpcreread.htm
share/doc/LPRng/lpdbounce.htm
share/doc/LPRng/lpf.htm
share/doc/LPRng/lpq.htm
share/doc/LPRng/lpr.htm
share/doc/LPRng/lprbounce.htm
share/doc/LPRng/lprm.htm
share/doc/LPRng/lprngclients.htm
share/doc/LPRng/lpsimulation.htm
share/doc/LPRng/maillist.htm
share/doc/LPRng/minfree.htm
share/doc/LPRng/nfsmount.htm
share/doc/LPRng/normalterm.htm
share/doc/LPRng/ofdetails.htm
share/doc/LPRng/opendevice.htm
share/doc/LPRng/optionindex.htm
share/doc/LPRng/p450.htm
share/doc/LPRng/permspath.htm
share/doc/LPRng/permsref.htm
share/doc/LPRng/pr.htm
share/doc/LPRng/printcappath.htm
share/doc/LPRng/printcapref.htm
share/doc/LPRng/printercomm.htm
share/doc/LPRng/printingandchecking.htm
share/doc/LPRng/printingjob.htm
share/doc/LPRng/printjobfilters.htm
share/doc/LPRng/printjobformats.htm
share/doc/LPRng/qq.htm
share/doc/LPRng/queuecontrolfile.htm
share/doc/LPRng/queuelockfile.htm
share/doc/LPRng/remotelogger.htm
share/doc/LPRng/remotesupport.htm
share/doc/LPRng/requireconfig.htm
share/doc/LPRng/restoring.htm
share/doc/LPRng/rfc1179.htm
share/doc/LPRng/rfc1179.txt
share/doc/LPRng/rfc1179ref.htm
share/doc/LPRng/rg.htm
share/doc/LPRng/safechars.htm
share/doc/LPRng/sampleprintcap.htm
share/doc/LPRng/secftp.htm
share/doc/LPRng/secnetwork.htm
share/doc/LPRng/secnis.htm
share/doc/LPRng/secoh.htm
share/doc/LPRng/secserial.htm
share/doc/LPRng/selectingprintqueue.htm
share/doc/LPRng/sendblockformat.htm
share/doc/LPRng/serial.htm
share/doc/LPRng/setuid.htm
share/doc/LPRng/simple.htm
share/doc/LPRng/smb.htm
share/doc/LPRng/socketapi.htm
share/doc/LPRng/spoolq.htm
share/doc/LPRng/spoolqueue.htm
share/doc/LPRng/stalledtime.htm
share/doc/LPRng/star1.gif
share/doc/LPRng/start.htm
share/doc/LPRng/submitting.htm
share/doc/LPRng/systemspecific.htm
share/doc/LPRng/testing.htm
share/doc/LPRng/topbanner.gif
share/doc/LPRng/translateformat.htm
share/doc/LPRng/tutorial.htm
share/doc/LPRng/tutorialconfig.htm
share/doc/LPRng/update.htm
share/doc/LPRng/updateversion
share/doc/LPRng/useidentifier.htm
share/doc/LPRng/usergroup.htm
share/doc/LPRng/usingcommandlineandprintcap.htm
share/doc/LPRng/x10158.htm
share/doc/LPRng/x10200.htm
share/doc/LPRng/x10205.htm
share/doc/LPRng/x10237.htm
share/doc/LPRng/x10247.htm
share/doc/LPRng/x10279.htm
share/doc/LPRng/x10292.htm
share/doc/LPRng/x10309.htm
share/doc/LPRng/x10334.htm
share/doc/LPRng/x10341.htm
share/doc/LPRng/x10351.htm
share/doc/LPRng/x10383.htm
share/doc/LPRng/x10394.htm
share/doc/LPRng/x10399.htm
share/doc/LPRng/x10427.htm
share/doc/LPRng/x10430.htm
share/doc/LPRng/x10435.htm
share/doc/LPRng/x10436.htm
share/doc/LPRng/x10441.htm
share/doc/LPRng/x10447.htm
share/doc/LPRng/x10469.htm
share/doc/LPRng/x10472.htm
share/doc/LPRng/x10477.htm
share/doc/LPRng/x10489.htm
share/doc/LPRng/x11153.htm
share/doc/LPRng/x11181.htm
share/doc/LPRng/x11195.htm
share/doc/LPRng/x11223.htm
share/doc/LPRng/x119.htm
share/doc/LPRng/x188.htm
share/doc/LPRng/x197.htm
share/doc/LPRng/x204.htm
share/doc/LPRng/x209.htm
share/doc/LPRng/x241.htm
share/doc/LPRng/x255.htm
share/doc/LPRng/x2691.htm
share/doc/LPRng/x2733.htm
share/doc/LPRng/x2751.htm
share/doc/LPRng/x2783.htm
share/doc/LPRng/x2793.htm
share/doc/LPRng/x2806.htm
share/doc/LPRng/x2825.htm
share/doc/LPRng/x2848.htm
share/doc/LPRng/x2881.htm
share/doc/LPRng/x2923.htm
share/doc/LPRng/x2980.htm
share/doc/LPRng/x3022.htm
share/doc/LPRng/x3039.htm
share/doc/LPRng/x3072.htm
share/doc/LPRng/x3081.htm
share/doc/LPRng/x3106.htm
share/doc/LPRng/x3114.htm
share/doc/LPRng/x3148.htm
share/doc/LPRng/x3209.htm
share/doc/LPRng/x3245.htm
share/doc/LPRng/x3251.htm
share/doc/LPRng/x3287.htm
share/doc/LPRng/x3292.htm
share/doc/LPRng/x3334.htm
share/doc/LPRng/x3337.htm
share/doc/LPRng/x3379.htm
share/doc/LPRng/x3415.htm
share/doc/LPRng/x3457.htm
share/doc/LPRng/x3516.htm
share/doc/LPRng/x3558.htm
share/doc/LPRng/x3589.htm
share/doc/LPRng/x3631.htm
share/doc/LPRng/x3707.htm
share/doc/LPRng/x3739.htm
share/doc/LPRng/x3749.htm
share/doc/LPRng/x3781.htm
share/doc/LPRng/x3931.htm
share/doc/LPRng/x3937.htm
share/doc/LPRng/x3943.htm
share/doc/LPRng/x395.htm
share/doc/LPRng/x3959.htm
share/doc/LPRng/x3966.htm
share/doc/LPRng/x3973.htm
share/doc/LPRng/x3979.htm
share/doc/LPRng/x3985.htm
share/doc/LPRng/x4001.htm
share/doc/LPRng/x4008.htm
share/doc/LPRng/x4083.htm
share/doc/LPRng/x4125.htm
share/doc/LPRng/x4189.htm
share/doc/LPRng/x4231.htm
share/doc/LPRng/x435.htm
share/doc/LPRng/x4381.htm
share/doc/LPRng/x4388.htm
share/doc/LPRng/x4399.htm
share/doc/LPRng/x4410.htm
share/doc/LPRng/x4423.htm
share/doc/LPRng/x4428.htm
share/doc/LPRng/x4430.htm
share/doc/LPRng/x4441.htm
share/doc/LPRng/x4452.htm
share/doc/LPRng/x447.htm
share/doc/LPRng/x4470.htm
share/doc/LPRng/x4528.htm
share/doc/LPRng/x4570.htm
share/doc/LPRng/x4577.htm
share/doc/LPRng/x4589.htm
share/doc/LPRng/x4598.htm
share/doc/LPRng/x46.htm
share/doc/LPRng/x4604.htm
share/doc/LPRng/x4619.htm
share/doc/LPRng/x4631.htm
share/doc/LPRng/x4640.htm
share/doc/LPRng/x4646.htm
share/doc/LPRng/x4657.htm
share/doc/LPRng/x4673.htm
share/doc/LPRng/x4699.htm
share/doc/LPRng/x4715.htm
share/doc/LPRng/x4782.htm
share/doc/LPRng/x4824.htm
share/doc/LPRng/x4963.htm
share/doc/LPRng/x4981.htm
share/doc/LPRng/x4995.htm
share/doc/LPRng/x5005.htm
share/doc/LPRng/x5019.htm
share/doc/LPRng/x5023.htm
share/doc/LPRng/x5037.htm
share/doc/LPRng/x5061.htm
share/doc/LPRng/x51.htm
share/doc/LPRng/x5485.htm
share/doc/LPRng/x5511.htm
share/doc/LPRng/x5523.htm
share/doc/LPRng/x5527.htm
share/doc/LPRng/x5540.htm
share/doc/LPRng/x5553.htm
share/doc/LPRng/x5565.htm
share/doc/LPRng/x5574.htm
share/doc/LPRng/x5582.htm
share/doc/LPRng/x5601.htm
share/doc/LPRng/x5616.htm
share/doc/LPRng/x5643.htm
share/doc/LPRng/x5870.htm
share/doc/LPRng/x5912.htm
share/doc/LPRng/x5983.htm
share/doc/LPRng/x5997.htm
share/doc/LPRng/x6018.htm
share/doc/LPRng/x6025.htm
share/doc/LPRng/x6039.htm
share/doc/LPRng/x6060.htm
share/doc/LPRng/x6415.htm
share/doc/LPRng/x6457.htm
share/doc/LPRng/x6559.htm
share/doc/LPRng/x6601.htm
share/doc/LPRng/x70.htm
share/doc/LPRng/x7121.htm
share/doc/LPRng/x7163.htm
share/doc/LPRng/x7643.htm
share/doc/LPRng/x7685.htm
share/doc/LPRng/x8042.htm
share/doc/LPRng/x8084.htm
share/doc/LPRng/x812.htm
share/doc/LPRng/x824.htm
share/doc/LPRng/x8312.htm
share/doc/LPRng/x8354.htm
share/doc/LPRng/x848.htm
share/doc/LPRng/x8551.htm
share/doc/LPRng/x8564.htm
share/doc/LPRng/x8579.htm
share/doc/LPRng/x8593.htm
share/doc/LPRng/x8594.htm
share/doc/LPRng/x860.htm
share/doc/LPRng/x8606.htm
share/doc/LPRng/x8612.htm
share/doc/LPRng/x8621.htm
share/doc/LPRng/x8636.htm
share/doc/LPRng/x8648.htm
share/doc/LPRng/x8654.htm
share/doc/LPRng/x8715.htm
share/doc/LPRng/x8747.htm
share/doc/LPRng/x8757.htm
share/doc/LPRng/x8789.htm
share/doc/LPRng/x882.htm
share/doc/LPRng/x894.htm
share/doc/LPRng/x9014.htm
share/doc/LPRng/x9040.htm
share/doc/LPRng/x9056.htm
share/doc/LPRng/x9082.htm
share/doc/LPRng/x9107.htm
share/doc/LPRng/x9149.htm
share/doc/LPRng/x9233.htm
share/doc/LPRng/x9275.htm
share/doc/LPRng/x9284.htm
share/doc/LPRng/x9295.htm
share/doc/LPRng/x9310.htm
share/doc/LPRng/x9326.htm
share/doc/LPRng/x9337.htm
share/doc/LPRng/x9352.htm
share/doc/LPRng/x9693.htm
share/doc/LPRng/x9729.htm
share/doc/LPRng/x9735.htm
share/doc/LPRng/x9771.htm
share/doc/LPRng/x9932.htm
share/doc/LPRng/x9955.htm
share/doc/LPRng/x9974.htm
share/doc/LPRng/x9997.htm
share/doc/LPRng/y2k.txt
@dirrm libexec/filters
@dirrm share/doc/LPRng
