# sysutils/LPRng/pkg-plist
bin/cancel
bin/lp
bin/lpq
bin/lprm
bin/lpr
bin/lpstat
sbin/lpc
sbin/lpd
libexec/filters/lpf
libexec/filters/lpbanner
libexec/filters/pclbanner
libexec/filters/psbanner
sbin/checkpc
sbin/lpraccnt
etc/rc.d/lprng.sh
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/lprng.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/Makefile
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/README
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/Macros/using.tmac.usenix
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/Macros/tmac.usenix
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/Macros/Makefile
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/bg.gif
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-1.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-LISA95.ps
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/rfc1179.txt
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/UPDATE
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO.txt
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/star1.gif
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO.ps
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-16.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO.dvi
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-9.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-8.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-7.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-6.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-5.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-4.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-3.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-21.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-20.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-2.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-19.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-18.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-17.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/REDIRECT.LPRng.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LISA98.ppt
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LISA98.ps
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/SENDTOWEBSITE
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-15.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-14.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-13.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-12.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-11.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/LPRng-HOWTO-10.html
share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/license.txt
@dirrm share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%/Macros
@dirrm share/doc/LPRng/LPRng_DOC-%%LPRNG_VER%%
@dirrm share/doc/LPRng
@dirrm libexec/filters
