# sysutils/LPRng/pkg-plist
bin/lpq
bin/lprm
bin/lpr
bin/lpc
bin/checkpc
bin/lpraccnt
etc/lpd.perms
etc/lpd.conf
etc/rc.d/lprng.sh
libexec/lpd
libexec/lpf
libexec/lpbanner
libexec/readfilecount
libexec/removeoneline
libexec/authenticate_pgp
libexec/setupauth
lib/liblpr.a
lib/liblpr.so.3.21
man/man1/cancel.1.gz
man/man1/lp.1.gz
man/man1/lpbanner.1.gz
man/man1/lpc.1.gz
man/man1/lpf.1.gz
man/man1/lpq.1.gz
man/man1/lpr.1.gz
man/man1/lpraccnt.1.gz
man/man1/lprm.1.gz
man/man1/lpstat.1.gz
man/man1/monitor.1.gz
man/man1/readfilecount.1.gz
man/man5/lpd.conf.5.gz
man/man5/lpd.perms.5.gz
man/man5/printcap.5.gz
man/man8/checkpc.8.gz
man/man8/lpd.8.gz
share/doc/LPRng/README
share/doc/LPRng/README.1st
share/doc/LPRng/README.LPRng_to_lpr_spooling
share/doc/LPRng/README.accounting
share/doc/LPRng/README.authentication
share/doc/LPRng/README.banner
share/doc/LPRng/README.bouncequeues
share/doc/LPRng/README.deskjet
share/doc/LPRng/README.install
share/doc/LPRng/README.kerberos_authentication
share/doc/LPRng/README.logger
share/doc/LPRng/README.lp_Malformed_from_address
share/doc/LPRng/README.lp_emulation
share/doc/LPRng/README.lpd.conf
share/doc/LPRng/README.lpd.perms
share/doc/LPRng/README.nis
share/doc/LPRng/README.pgp_authentication
share/doc/LPRng/README.ports
share/doc/LPRng/README.printcap
share/doc/LPRng/README.routing
share/doc/LPRng/README.solaris
share/doc/LPRng/README.solaris.newsprint
share/doc/LPRng/README.statusdisplays
@dirrm share/doc/LPRng
