# sysutils/sysadm/pkg-plist
bin/sysadm
bin/sysadm-server
etc/rc.d/sysadm-restserver
etc/rc.d/sysadm-websocket
include/sysadm-firewall.h
include/sysadm-general.h
include/sysadm-global.h
include/sysadm-iocage.h
include/sysadm-lifepreserver.h
include/sysadm-network.h
include/sysadm-servicemanager.h
include/sysadm-systeminfo.h
include/sysadm-update.h
include/sysadm-usermanager.h
%%QT_LIBDIR%%/%%QT_LIBDIR%%sysadm.so
%%QT_LIBDIR%%/%%QT_LIBDIR%%sysadm.so.1
%%QT_LIBDIR%%/%%QT_LIBDIR%%sysadm.so.1.0
%%QT_LIBDIR%%/%%QT_LIBDIR%%sysadm.so.1.0.0
