# sysutils/sysadm/pkg-plist
bin/sysadm-binary
bin/sysadm-server
bin/sysadm-bridge
bin/sysadm-bridge-start
etc/launchd/daemons/org.pcbsd.sysadm.json.dist
etc/launchd/daemons/org.pcbsd.sysadm-rest.json.dist
etc/sysadm.conf.dist
etc/rc.d/sysadm
etc/rc.d/sysadm-bridge
etc/rc.d/sysadm-rest
share/launchd/daemons/org.pcbsd.sysadm.json
share/launchd/daemons/org.pcbsd.sysadm-rest.json
