# sysutils/sysadm/pkg-plist
bin/sysadm-server
etc/sysadm.conf.dist
etc/rc.d/sysadm
etc/rc.d/sysadm-rest
