# net-mgmt/unifi5/files/openrc-unifi.in
#!/sbin/openrc-run
#
# Original Author: Mark Felder <feld@FreeBSD.org>
# Author:  Mike N. <c0psrul3@gmail.com>
# $FreeBSD$
#

# PROVIDE: unifi
# REQUIRE: LOGIN
# KEYWORD: shutdown

name=unifi

#rcvar=unifi_enable

#: ${unifi_enable:=NO}
: ${unifi_user:=unifi}
: ${unifi_group:=unifi}
: ${unifi_chdir=%%PREFIX%%/share/java/${name}}
: ${unifi_javaflags:="-Djava.awt.headless=true -Xmx1024M"}

piddir="/var/run/${name}"
pidfile="${piddir}/${name}.pid"

if [ -x %%PREFIX%%/openjdk8-jre/bin/java ] ; then
  JREbin='%%PREFIX%%/openjdk8-jre/bin/java'
else
  for _J in $(find %%PREFIX%% -type d -path "%%PREFIX%%/openjdk*" -name "bin" 2> /dev/null) ; do
    if [ -x ${_J}/java ] ; then
      Jbin="${_J}/java"
      JREbin=$(echo ${Jbin} | grep "jre")
    fi
  done
fi
procname=${JREbin:-$Jbin}

command_args="${procname} ${unifi_javaflags} -jar lib/ace.jar start"

start() {
    ebegin "Starting ${name}"
    start-stop-daemon --start --progress --make-pidfile --background \
        --user ${unifi_user}:${unifi_group} \
        --chdir ${unifi_chdir} \
        --pidfile ${pidfile} \
        -- ${procname} ${unifi_javaflags} -jar lib/ace.jar start
    eend $?
}

stop() {
    ebegin "Stopping ${name}"
    start-stop-daemon --stop --exec ${procname} \
        --pidfile ${pidfile}
    eend $?
}

start_pre()
{
	if [ ! -e  ${piddir} ] ; then
		install -d -o ${unifi_user} -g ${unifi_group} ${piddir};
	fi
}
