# net-mgmt/unifi5/files/openrc-unifi.in
#!/sbin/openrc-run
#
# Original Author: Mark Felder <feld@FreeBSD.org>
# Author:  Mike N. <c0psrul3@gmail.com>

name=unifi

: ${unifi_user:=unifi}
: ${unifi_group:=unifi}
: ${unifi_chdir=%%PREFIX%%/share/java/${name}}
: ${unifi_javaflags:="-Djava.awt.headless=true -Xmx1024M"}

piddir="/var/run/${name}"
pidfile="${piddir}/${name}.pid"

if [ -x %%PREFIX%%/openjdk8-jre/bin/java ] ; then
  JREbin='%%PREFIX%%/openjdk8-jre/bin/java'
else
  for _J in $(find %%PREFIX%% -type d -path "%%PREFIX%%/openjdk*" -name "bin" 2> /dev/null) ; do
    if [ -x ${_J}/java ] ; then
      Jbin="${_J}/java"
      JREbin=$(echo ${Jbin} | grep "jre")
    fi
  done
fi

command=${JREbin:-$Jbin}
command_args="${unifi_javaflags} -jar lib/ace.jar start"
command_user=${unifi_user}

start_pre() {
	if [ ! -e  ${piddir} ] ; then
		install -d -o ${unifi_user} -g ${unifi_group} ${piddir};
	fi
}

depend() {
	need localmount
	keyword -shutdown
	after bootmisc
}
