# www/sams2/files/sams2
#!/bin/sh

# $FreeBSD

# PROVIDE: sams2
# REQUIRE: LOGIN cleanvar mysql
# KEYWORD: shutdown

#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#
# sams2_enable="YES"
# sams2_flags="<set as needed>"

sams2_enable=${sams2_enable-"NO"}
sams2_pidfile=${sams2_pidfile-"/var/run/sams2daemon.pid"}

. /etc/rc.subr

name="sams2"
rcvar=sams2_enable
command="/usr/local/bin/sams2daemon"

load_rc_config $name

pidfile="${sams2_pidfile}"

run_rc_command "$1"

