# www/aolserver/pkg-plist
aolserver/bin/translate-ini
aolserver/bin/nsd8x
aolserver/bin/nsd76
aolserver/bin/nsd
aolserver/bin/nssock.so
aolserver/bin/nscgi.so
aolserver/bin/nscp.so
aolserver/bin/nslog.so
aolserver/bin/nsperm.so
aolserver/bin/nsext.so
aolserver/bin/nsunix.so
aolserver/bin/nsvhr.so
aolserver/modules/tcl/util.tcl
aolserver/modules/tcl/debug.tcl
aolserver/modules/tcl/namespace.tcl
aolserver/modules/tcl/file.tcl
aolserver/modules/tcl/form.tcl
aolserver/modules/tcl/http.tcl
aolserver/modules/tcl/init.tcl
aolserver/modules/tcl/modlog.tcl
aolserver/modules/tcl/prodebug.tcl
aolserver/modules/tcl/nsdb.tcl
aolserver/modules/tcl/sendmail.tcl
aolserver/modules/tcl/tcl.html
aolserver/modules/tcl/fastpath.tcl
aolserver/modules/tcl/nsperm/compat.tcl
aolserver/modules/tcl/nsperm/init.tcl
aolserver/lib/tcl8.3/msgcat1.0/pkgIndex.tcl
aolserver/lib/tcl8.3/msgcat1.0/msgcat.tcl
aolserver/lib/tcl8.3/encoding/symbol.enc
aolserver/lib/tcl8.3/encoding/ascii.enc
aolserver/lib/tcl8.3/encoding/big5.enc
aolserver/lib/tcl8.3/encoding/cp1250.enc
aolserver/lib/tcl8.3/encoding/cp1251.enc
aolserver/lib/tcl8.3/encoding/cp1252.enc
aolserver/lib/tcl8.3/encoding/cp1253.enc
aolserver/lib/tcl8.3/encoding/cp1254.enc
aolserver/lib/tcl8.3/encoding/cp1255.enc
aolserver/lib/tcl8.3/encoding/cp1256.enc
aolserver/lib/tcl8.3/encoding/cp1257.enc
aolserver/lib/tcl8.3/encoding/cp1258.enc
aolserver/lib/tcl8.3/encoding/cp437.enc
aolserver/lib/tcl8.3/encoding/cp737.enc
aolserver/lib/tcl8.3/encoding/cp775.enc
aolserver/lib/tcl8.3/encoding/cp850.enc
aolserver/lib/tcl8.3/encoding/cp852.enc
aolserver/lib/tcl8.3/encoding/cp855.enc
aolserver/lib/tcl8.3/encoding/cp857.enc
aolserver/lib/tcl8.3/encoding/cp860.enc
aolserver/lib/tcl8.3/encoding/cp861.enc
aolserver/lib/tcl8.3/encoding/cp862.enc
aolserver/lib/tcl8.3/encoding/cp863.enc
aolserver/lib/tcl8.3/encoding/cp864.enc
aolserver/lib/tcl8.3/encoding/cp865.enc
aolserver/lib/tcl8.3/encoding/cp866.enc
aolserver/lib/tcl8.3/encoding/cp869.enc
aolserver/lib/tcl8.3/encoding/cp874.enc
aolserver/lib/tcl8.3/encoding/cp932.enc
aolserver/lib/tcl8.3/encoding/cp936.enc
aolserver/lib/tcl8.3/encoding/cp949.enc
aolserver/lib/tcl8.3/encoding/cp950.enc
aolserver/lib/tcl8.3/encoding/gb12345.enc
aolserver/lib/tcl8.3/encoding/euc-cn.enc
aolserver/lib/tcl8.3/encoding/euc-jp.enc
aolserver/lib/tcl8.3/encoding/euc-kr.enc
aolserver/lib/tcl8.3/encoding/iso2022-jp.enc
aolserver/lib/tcl8.3/encoding/gb1988.enc
aolserver/lib/tcl8.3/encoding/gb2312.enc
aolserver/lib/tcl8.3/encoding/macCentEuro.enc
aolserver/lib/tcl8.3/encoding/iso2022-kr.enc
aolserver/lib/tcl8.3/encoding/iso2022.enc
aolserver/lib/tcl8.3/encoding/iso8859-1.enc
aolserver/lib/tcl8.3/encoding/iso8859-2.enc
aolserver/lib/tcl8.3/encoding/iso8859-3.enc
aolserver/lib/tcl8.3/encoding/iso8859-4.enc
aolserver/lib/tcl8.3/encoding/shiftjis.enc
aolserver/lib/tcl8.3/encoding/jis0201.enc
aolserver/lib/tcl8.3/encoding/iso8859-5.enc
aolserver/lib/tcl8.3/encoding/iso8859-6.enc
aolserver/lib/tcl8.3/encoding/iso8859-7.enc
aolserver/lib/tcl8.3/encoding/iso8859-8.enc
aolserver/lib/tcl8.3/encoding/iso8859-9.enc
aolserver/lib/tcl8.3/encoding/jis0208.enc
aolserver/lib/tcl8.3/encoding/jis0212.enc
aolserver/lib/tcl8.3/encoding/koi8-r.enc
aolserver/lib/tcl8.3/encoding/ksc5601.enc
aolserver/lib/tcl8.3/encoding/macCroatian.enc
aolserver/lib/tcl8.3/encoding/macCyrillic.enc
aolserver/lib/tcl8.3/encoding/macDingbats.enc
aolserver/lib/tcl8.3/encoding/macGreek.enc
aolserver/lib/tcl8.3/encoding/macIceland.enc
aolserver/lib/tcl8.3/encoding/macJapan.enc
aolserver/lib/tcl8.3/encoding/macRoman.enc
aolserver/lib/tcl8.3/encoding/macRomania.enc
aolserver/lib/tcl8.3/encoding/macThai.enc
aolserver/lib/tcl8.3/encoding/macTurkish.enc
aolserver/lib/tcl8.3/encoding/macUkraine.enc
aolserver/lib/tcl8.3/encoding/dingbats.enc
aolserver/lib/tcl8.3/http1.0/pkgIndex.tcl
aolserver/lib/tcl8.3/http1.0/http.tcl
aolserver/lib/tcl8.3/http2.3/pkgIndex.tcl
aolserver/lib/tcl8.3/http2.3/http.tcl
aolserver/lib/tcl8.3/dde1.1/pkgIndex.tcl
aolserver/lib/tcl8.3/opt0.4/optparse.tcl
aolserver/lib/tcl8.3/opt0.4/pkgIndex.tcl
aolserver/lib/tcl8.3/reg1.0/pkgIndex.tcl
aolserver/lib/tcl8.3/tcltest1.0/pkgIndex.tcl
aolserver/lib/tcl8.3/tcltest1.0/tcltest.tcl
aolserver/lib/tcl8.3/history.tcl
aolserver/lib/tcl8.3/word.tcl
aolserver/lib/tcl8.3/tclIndex
aolserver/lib/tcl8.3/auto.tcl
aolserver/lib/tcl8.3/license.terms
aolserver/lib/tcl8.3/init.tcl
aolserver/lib/tcl8.3/ldAout.tcl
aolserver/lib/tcl8.3/package.tcl
aolserver/lib/tcl8.3/parray.tcl
aolserver/lib/tcl8.3/safe.tcl
aolserver/lib/libnspd.a
aolserver/lib/libnsthread.a
aolserver/servers/server1/modules/nsperm/perms
aolserver/servers/server1/modules/nsperm/hosts.deny
aolserver/servers/server1/modules/nsperm/hosts.allow
aolserver/servers/server1/modules/nsperm/group
aolserver/servers/server1/modules/nsperm/passwd
aolserver/servers/server1/pages/index.html
aolserver/include/tclDecls.h
aolserver/include/Makefile.module
aolserver/include/ns.h
aolserver/include/nsextmsg.h
aolserver/include/nspd.h
aolserver/include/nsthread.h
aolserver/include/tcl.h
aolserver/include/tcl76.h
aolserver/include/tcl83.h
aolserver/include/Makefile.global
aolserver/include/nsd.h
aolserver/sample-config.tcl
etc/rc.d/aolserver.sh
share/doc/aolserver/ChangeLog
share/doc/aolserver/README
share/doc/aolserver/install.html
share/doc/aolserver/config.txt
share/doc/aolserver/faq.html
share/doc/aolserver/default-home.html
share/doc/aolserver/tcl-changes.txt
share/doc/aolserver/license.txt
share/doc/aolserver/tuning.html
@dirrm aolserver/bin
@dirrm aolserver/log
@dirrm aolserver/modules/tcl/nsperm
@dirrm aolserver/modules/tcl
@dirrm aolserver/modules
@dirrm aolserver/lib/tcl8.3/msgcat1.0
@dirrm aolserver/lib/tcl8.3/encoding
@dirrm aolserver/lib/tcl8.3/http1.0
@dirrm aolserver/lib/tcl8.3/http2.3
@dirrm aolserver/lib/tcl8.3/dde1.1
@dirrm aolserver/lib/tcl8.3/opt0.4
@dirrm aolserver/lib/tcl8.3/reg1.0
@dirrm aolserver/lib/tcl8.3/tcltest1.0
@dirrm aolserver/lib/tcl8.3
@dirrm aolserver/lib
@dirrm aolserver/servers/server1/modules/nslog
@dirrm aolserver/servers/server1/modules/nsperm
@dirrm aolserver/servers/server1/modules
@dirrm aolserver/servers/server1/pages
@dirrm aolserver/servers/server1
@dirrm aolserver/servers
@dirrm aolserver/include
@dirrm aolserver
@dirrm share/doc/aolserver
@dirrm share/doc
