# www/aolserver/pkg-plist
aolserver/bin/translate-ini
aolserver/bin/nsd8x
aolserver/bin/nsd76
aolserver/bin/nsd
aolserver/bin/nssock.so
aolserver/bin/nscgi.so
aolserver/bin/nscp.so
aolserver/bin/nslog.so
aolserver/bin/nsperm.so
aolserver/bin/nsext.so
aolserver/bin/nsunix.so
aolserver/bin/nsvhr.so
aolserver/modules/tcl/util.tcl
aolserver/modules/tcl/fastpath.tcl
aolserver/modules/tcl/file.tcl
aolserver/modules/tcl/form.tcl
aolserver/modules/tcl/http.tcl
aolserver/modules/tcl/init.tcl
aolserver/modules/tcl/modlog.tcl
aolserver/modules/tcl/namespace.tcl
aolserver/modules/tcl/nsdb.tcl
aolserver/modules/tcl/prodebug.tcl
aolserver/modules/tcl/sendmail.tcl
aolserver/modules/tcl/debug.tcl
aolserver/modules/tcl/nsperm/compat.tcl
aolserver/modules/tcl/nsperm/init.tcl
aolserver/lib/libnspd.a
aolserver/lib/libnsthread.a
aolserver/servers/server1/modules/nsperm/perms
aolserver/servers/server1/modules/nsperm/hosts.deny
aolserver/servers/server1/modules/nsperm/hosts.allow
aolserver/servers/server1/modules/nsperm/group
aolserver/servers/server1/modules/nsperm/passwd
aolserver/servers/server1/pages/index.html
aolserver/include/tclDecls.h
aolserver/include/Makefile.module
aolserver/include/ns.h
aolserver/include/nsextmsg.h
aolserver/include/nspd.h
aolserver/include/nsthread.h
aolserver/include/tcl.h
aolserver/include/tcl76.h
aolserver/include/tcl83.h
aolserver/include/Makefile.global
aolserver/nsd.tcl.sample
etc/rc.d/aolserver.sh
include/aolserver/ns.h
include/aolserver/nsextmsg.h
include/aolserver/nspd.h
include/aolserver/nsthread.h
include/aolserver/tcl.h
include/aolserver/tcl76.h
include/aolserver/tcl83.h
include/aolserver/tclDecls.h
include/aolserver/nsd.h
lib/libnspd.a
share/doc/aolserver/ChangeLog
share/doc/aolserver/LICENSE.txt
share/doc/aolserver/README.txt
share/doc/aolserver/doc/config.txt
share/doc/aolserver/doc/install.txt
share/doc/aolserver/doc/internals.txt
share/doc/aolserver/doc/nsv.txt
share/doc/aolserver/doc/threads.txt
@dirrm aolserver/bin
@dirrm aolserver/log
@dirrm aolserver/modules/tcl/nsperm
@dirrm aolserver/modules/tcl
@dirrm aolserver/modules
@dirrm aolserver/lib
@dirrm aolserver/servers/server1/modules/nslog
@dirrm aolserver/servers/server1/modules/nsperm
@dirrm aolserver/servers/server1/modules
@dirrm aolserver/servers/server1/pages
@dirrm aolserver/servers/server1
@dirrm aolserver/servers
@dirrm aolserver/include
@dirrm aolserver
@dirrm include/aolserver
@dirrm include
@dirrm share/doc/aolserver/doc
@dirrm share/doc/aolserver
@dirrm share/doc
