# databases/hk_classes/pkg-plist
etc/rc.d/hk_classes.sh
bin/hk_actionquery
bin/hk_exportcsv
bin/hk_exporthtml
bin/hk_exportxml
bin/hk_importcsv
bin/hk_report
include/hk_classes/hk_actionquery.h
include/hk_classes/hk_button.h
include/hk_classes/hk_class.h
include/hk_classes/hk_classes.h
include/hk_classes/hk_colour.h
include/hk_classes/hk_column.h
include/hk_classes/hk_connection.h
include/hk_classes/hk_data.h
include/hk_classes/hk_database.h
include/hk_classes/hk_datasource.h
include/hk_classes/hk_datetime.h
include/hk_classes/hk_dbvisible.h
include/hk_classes/hk_definitions.h
include/hk_classes/hk_drivermanager.h
include/hk_classes/hk_dsboolean.h
include/hk_classes/hk_dscombobox.h
include/hk_classes/hk_dsdatavisible.h
include/hk_classes/hk_dsdate.h
include/hk_classes/hk_dsgrid.h
include/hk_classes/hk_dsgridcolumn.h
include/hk_classes/hk_dsimage.h
include/hk_classes/hk_dslineedit.h
include/hk_classes/hk_dsmemo.h
include/hk_classes/hk_dsmodevisible.h
include/hk_classes/hk_dsquery.h
include/hk_classes/hk_dsrowselector.h
include/hk_classes/hk_dstable.h
include/hk_classes/hk_dsvisible.h
include/hk_classes/hk_font.h
include/hk_classes/hk_form.h
include/hk_classes/hk_importcsv.h
include/hk_classes/hk_interpreter.h
include/hk_classes/hk_key.h
include/hk_classes/hk_label.h
include/hk_classes/hk_presentation.h
include/hk_classes/hk_qbe.h
include/hk_classes/hk_report.h
include/hk_classes/hk_reportdata.h
include/hk_classes/hk_reportsection.h
include/hk_classes/hk_reportsectionpair.h
include/hk_classes/hk_reportutils.h
include/hk_classes/hk_storagecolumn.h
include/hk_classes/hk_storagedatasource.h
include/hk_classes/hk_string.h
include/hk_classes/hk_subform.h
include/hk_classes/hk_tabvisible.h
include/hk_classes/hk_url.h
include/hk_classes/hk_visible.h
%%MYSQL%%lib/hk_classes/drivers/libhk_mysqldriver.so
%%MYSQL%%lib/hk_classes/drivers/libhk_mysqldriver.so.3
%%MYSQL%%lib/hk_classes/drivers/libhk_mysqldriver.a
%%ODBC%%lib/hk_classes/drivers/libhk_odbcdriver.so
%%ODBC%%lib/hk_classes/drivers/libhk_odbcdriver.so.0
%%ODBC%%lib/hk_classes/drivers/libhk_odbcdriver.a
%%POSTGRESQL%%lib/hk_classes/drivers/libhk_postgresdriver.so
%%POSTGRESQL%%lib/hk_classes/drivers/libhk_postgresdriver.so.0
%%POSTGRESQL%%lib/hk_classes/drivers/libhk_postgresdriver.a
%%SQLITE2%%lib/hk_classes/drivers/libhk_sqlite2driver.so
%%SQLITE2%%lib/hk_classes/drivers/libhk_sqlite2driver.so.0
%%SQLITE2%%lib/hk_classes/drivers/libhk_sqlite2driver.a
%%SQLITE3%%lib/hk_classes/drivers/libhk_sqlite3driver.so
%%SQLITE3%%lib/hk_classes/drivers/libhk_sqlite3driver.so.0
%%SQLITE3%%lib/hk_classes/drivers/libhk_sqlite3driver.a
%%PYTHON_SITELIBDIR%%/_hk_classes.so
%%PYTHON_SITELIBDIR%%/hk_classes.py
%%PYTHON_SITELIBDIR%%/hk_classes.pyc
lib/hk_classes/libhk_classes.so
lib/hk_classes/libhk_classes.so.8
lib/hk_classes/libhk_classes.a
@dirrm lib/hk_classes/drivers
@dirrm lib/hk_classes
@dirrm include/hk_classes
