# net/ssltunnel-server/files/openrc-ssltunnel.in
#!/sbin/openrc-run

name=ssltunnel
command="%%PREFIX%%/libexec/pppserver"
command_args="%%ETCDIR%%/tunnel.conf"
required_files="%%ETCDIR%%/tunnel.conf"
pidfile=/var/run/pppserver.pid

depend() {
	keyword -shutdown
}
