# net/ssltunnel-server/files/ssltunnel.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: ssltunnel
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# To enable ssltunnel, add the following line to rc.conf:
#	ssltunnel_enable=yes

. /etc/rc.subr

name=ssltunnel
rcvar=ssltunnel_enable

load_rc_config ssltunnel

: ${ssltunnel_enable:=yes}

command="%%PREFIX%%/libexec/pppserver"
command_args="%%ETCDIR%%/tunnel.conf"
required_files="%%ETCDIR%%/tunnel.conf"
pidfile=/var/run/pppserver.pid

run_rc_command $1
