# databases/mysql41-server/pkg-plist
bin/comp_err
bin/isamchk
bin/isamlog
bin/make_win_src_distribution
bin/msql2mysql
bin/my_print_defaults
bin/myisamchk
bin/myisamlog
bin/myisampack
bin/mysql_config
bin/mysql_convert_table_format
bin/mysql_create_system_tables
bin/mysql_explain_log
bin/mysql_find_rows
bin/mysql_fix_extensions
bin/mysql_fix_privilege_tables
bin/mysql_install_db
bin/mysql_secure_installation
bin/mysql_setpermission
bin/mysql_tableinfo
bin/mysql_waitpid
bin/mysql_zap
bin/mysqlaccess
bin/mysqlbug
bin/mysqld_multi
bin/mysqld_safe
bin/mysqldumpslow
bin/mysqlhotcopy
bin/pack_isam
bin/perror
bin/replace
bin/resolve_stack_dump
bin/resolveip
etc/rc.d/mysql-server.sh
lib/mysql/libdbug.a
lib/mysql/libheap.a
lib/mysql/libmerge.a
lib/mysql/libmyisam.a
lib/mysql/libmyisammrg.a
lib/mysql/libmystrings.a
lib/mysql/libmysys.a
lib/mysql/libnisam.a
lib/mysql/libvio.a
libexec/mysqld
%%PORTDOCS%%%%DOCSDIR%%/Flags/argentina.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/armenia.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/australia.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/austria.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/belgium.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/brazil.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/bulgaria.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/canada.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/chile.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/china.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/costa-rica.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/croatia.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/czech-republic.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/denmark.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/estonia.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/finland.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/france.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/germany.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/great-britain.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/greece.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/hungary.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/iceland.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/indonesia.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/ireland.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/israel.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/italy.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/japan.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/latvia.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/mexico.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/netherlands.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/new-zealand.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/norway.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/philippines.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/poland.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/portugal.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/romania.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/russia.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/singapore.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/slovenia.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/south-africa.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/south-korea.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/spain.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/sweden.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/switzerland.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/taiwan.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/turkey.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/ukraine.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/usa.gif
%%PORTDOCS%%%%DOCSDIR%%/Flags/yugoslavia.gif
%%PORTDOCS%%%%DOCSDIR%%/manual.html
%%PORTDOCS%%%%DOCSDIR%%/manual.txt
%%PORTDOCS%%%%DOCSDIR%%/manual_toc.html
share/mysql/Description.plist
share/mysql/Info.plist
share/mysql/MySQL-shared-compat.spec
share/mysql/StartupParameters.plist
share/mysql/postinstall
share/mysql/preinstall
share/mysql/charsets/Index.xml
share/mysql/charsets/README
share/mysql/charsets/armscii8.xml
share/mysql/charsets/ascii.xml
share/mysql/charsets/cp1250.xml
share/mysql/charsets/cp1251.xml
share/mysql/charsets/cp1256.xml
share/mysql/charsets/cp1257.xml
share/mysql/charsets/cp850.xml
share/mysql/charsets/cp852.xml
share/mysql/charsets/cp866.xml
share/mysql/charsets/dec8.xml
share/mysql/charsets/geostd8.xml
share/mysql/charsets/greek.xml
share/mysql/charsets/hebrew.xml
share/mysql/charsets/hp8.xml
share/mysql/charsets/keybcs2.xml
share/mysql/charsets/koi8r.xml
share/mysql/charsets/koi8u.xml
share/mysql/charsets/latin1.xml
share/mysql/charsets/latin2.xml
share/mysql/charsets/latin5.xml
share/mysql/charsets/latin7.xml
share/mysql/charsets/macce.xml
share/mysql/charsets/macroman.xml
share/mysql/charsets/swe7.xml
share/mysql/czech/errmsg.sys
share/mysql/czech/errmsg.txt
share/mysql/danish/errmsg.sys
share/mysql/danish/errmsg.txt
share/mysql/dutch/errmsg.sys
share/mysql/dutch/errmsg.txt
share/mysql/english/errmsg.sys
share/mysql/english/errmsg.txt
share/mysql/estonian/errmsg.sys
share/mysql/estonian/errmsg.txt
share/mysql/fill_help_tables.sql
share/mysql/french/errmsg.sys
share/mysql/french/errmsg.txt
share/mysql/german/errmsg.sys
share/mysql/german/errmsg.txt
share/mysql/greek/errmsg.sys
share/mysql/greek/errmsg.txt
share/mysql/hungarian/errmsg.sys
share/mysql/hungarian/errmsg.txt
share/mysql/italian/errmsg.sys
share/mysql/italian/errmsg.txt
share/mysql/japanese/errmsg.sys
share/mysql/japanese/errmsg.txt
share/mysql/korean/errmsg.sys
share/mysql/korean/errmsg.txt
share/mysql/mi_test_all
share/mysql/mi_test_all.res
share/mysql/my-huge.cnf
share/mysql/my-large.cnf
share/mysql/my-medium.cnf
share/mysql/my-small.cnf
share/mysql/mysql-log-rotate
share/mysql/mysql_fix_privilege_tables.sql
share/mysql/mysql.server
share/mysql/norwegian-ny/errmsg.sys
share/mysql/norwegian-ny/errmsg.txt
share/mysql/norwegian/errmsg.sys
share/mysql/norwegian/errmsg.txt
share/mysql/polish/errmsg.sys
share/mysql/polish/errmsg.txt
share/mysql/portuguese/errmsg.sys
share/mysql/portuguese/errmsg.txt
share/mysql/romanian/errmsg.sys
share/mysql/romanian/errmsg.txt
share/mysql/russian/errmsg.sys
share/mysql/russian/errmsg.txt
share/mysql/slovak/errmsg.sys
share/mysql/slovak/errmsg.txt
share/mysql/spanish/errmsg.sys
share/mysql/spanish/errmsg.txt
share/mysql/swedish/errmsg.sys
share/mysql/swedish/errmsg.txt
share/mysql/ukrainian/errmsg.sys
share/mysql/ukrainian/errmsg.txt
%%PORTDOCS%%@dirrm %%DOCSDIR%%/Flags
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm share/mysql/charsets
@dirrm share/mysql/czech
@dirrm share/mysql/danish
@dirrm share/mysql/dutch
@dirrm share/mysql/english
@dirrm share/mysql/estonian
@dirrm share/mysql/french
@dirrm share/mysql/german
@dirrm share/mysql/greek
@dirrm share/mysql/hungarian
@dirrm share/mysql/italian
@dirrm share/mysql/japanese
@dirrm share/mysql/korean
@dirrm share/mysql/norwegian
@dirrm share/mysql/norwegian-ny
@dirrm share/mysql/polish
@dirrm share/mysql/portuguese
@dirrm share/mysql/romanian
@dirrm share/mysql/russian
@dirrm share/mysql/slovak
@dirrm share/mysql/spanish
@dirrm share/mysql/swedish
@dirrm share/mysql/ukrainian
@dirrm share/mysql
@exec [ -d %%DB_DIR%% ] || %D/bin/mysql_install_db --ldata=%%DB_DIR%%
@unexec /usr/bin/killall mysqld > /dev/null 2>&1 || true
