# sysutils/puppet37/files/puppet.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: puppet
# REQUIRE: NETWORK
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable the puppet agent:
#
# puppet_enable="YES"

. /etc/rc.subr

name="puppet"
rcvar=puppet_enable

load_rc_config "$name"

: ${puppet_enable="NO"}
: ${puppet_rundir="/var/run/puppet"}

command="%%PREFIX%%/bin/puppet"
command_args="agent ${puppet_flags} --rundir=${puppet_rundir}"
command_interpreter=%%RUBY%%
unset puppet_flags

pidfile="${puppet_rundir}/agent.pid"
start_precmd="install -d -o puppet -g puppet ${pidfile%/*}"

run_rc_command "$1"
