# sysutils/bsdstats/pkg-plist
etc/periodic/monthly/300.statistics
etc/rc.d/bsdstats.sh
@dirrmtry etc/periodic/monthly
@dirrmtry etc/periodic
