# misc/trueos-server/pkg-plist
bin/disable-tor-mode
bin/enable-tor-mode
bin/migrate_rc_openrc
bin/restamp-grub
bin/trueos-mount
bin/trueos-unmount
bin/warden-to-iocage
etc/devd/automount.conf
etc/fonts/local.conf
etc/sudoers.d/trueos
share/fluxbox/styles/bora_blue_4k/theme.cfg
share/icons/default/index.theme
share/pixmaps/start-here-lumina.png
share/trueos/cardDetect/XF86Config.compat
share/trueos/cardDetect/XF86Config.default
share/trueos/cardDetect/XF86Config.intel
share/trueos/cardDetect/XF86Config.modesetting
share/trueos/cardDetect/XF86Config.qemu
share/trueos/cardDetect/XF86Config.scfb
share/trueos/cardDetect/XF86Config.virtualbox
share/trueos/cardDetect/XF86Config.vmware
share/trueos/cardDetect/setuptouch.sh
share/trueos/cardDetect/x_config.sh
share/trueos/cardDetect/xorg.conf.virtualbox
share/trueos/conf/brand-trueos.4th
share/trueos/conf/logo-trueos.4th
share/trueos/conf/ntpd.conf
share/trueos/conf/rc.conf.desktop
share/trueos/conf/rc.conf.server
share/trueos/conf/smb4.conf
share/trueos/conf/torrc
share/trueos/desktop-defaults/boot/grub/themes/trueos/README
share/trueos/desktop-defaults/boot/grub/themes/trueos/background.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/blob_w.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/boot_menu_c.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/boot_menu_e.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/boot_menu_n.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/boot_menu_ne.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/boot_menu_nw.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/boot_menu_s.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/boot_menu_se.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/boot_menu_sw.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/boot_menu_w.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/dejavu_10.pf2
share/trueos/desktop-defaults/boot/grub/themes/trueos/dejavu_12.pf2
share/trueos/desktop-defaults/boot/grub/themes/trueos/dejavu_14.pf2
share/trueos/desktop-defaults/boot/grub/themes/trueos/dejavu_16.pf2
share/trueos/desktop-defaults/boot/grub/themes/trueos/dejavu_bold_14.pf2
share/trueos/desktop-defaults/boot/grub/themes/trueos/slider_c.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/slider_n.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/slider_s.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/terminal_box_c.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/terminal_box_e.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/terminal_box_n.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/terminal_box_ne.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/terminal_box_nw.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/terminal_box_s.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/terminal_box_se.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/terminal_box_sw.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/terminal_box_w.png
share/trueos/desktop-defaults/boot/grub/themes/trueos/theme.txt
share/trueos/desktop-defaults/boot/grub/trueosfont.pf2
share/trueos/desktop-defaults/etc/blacklist
share/trueos/desktop-defaults/etc/crontab
share/trueos/desktop-defaults/etc/devfs.conf
share/trueos/desktop-defaults/etc/devfs.rules
share/trueos/desktop-defaults/etc/hosts.allow
share/trueos/desktop-defaults/etc/hosts.deniedssh
share/trueos/desktop-defaults/etc/timezone
share/trueos/desktop-defaults/etc/wall_cmos_clock
share/trueos/desktop-defaults/root/.bashrc
share/trueos/desktop-defaults/root/.config/qtcurve/gtk-icons
share/trueos/desktop-defaults/root/.config/qtcurve/stylerc
share/trueos/desktop-defaults/root/.cshrc
share/trueos/desktop-defaults/root/.kde4/share/config/kdeglobals
share/trueos/desktop-defaults/root/.kde4/share/config/kwinrc
share/trueos/desktop-defaults/root/.kde4/share/config/nitrogenrc
share/trueos/desktop-defaults/root/.profile
share/trueos/desktop-defaults/root/.qt/lipstikstylerc
share/trueos/desktop-defaults/root/.vimrc
share/trueos/desktop-defaults/usr/local/bin/flashpluginctl
share/trueos/desktop-defaults/usr/local/bin/openwith
share/trueos/desktop-defaults/usr/local/bin/trueos-update-i18n
share/trueos/desktop-defaults/usr/local/bin/xlock
share/trueos/desktop-defaults/usr/local/bin/xorg-startupcheck
share/trueos/desktop-defaults/usr/local/etc/beadm.conf
share/trueos/desktop-defaults/usr/local/etc/cups/cups-pdf.conf
share/trueos/desktop-defaults/usr/local/etc/cups/cupsd.conf
share/trueos/desktop-defaults/usr/local/etc/default/grub
share/trueos/desktop-defaults/usr/local/etc/foomatic/filter.conf
share/trueos/desktop-defaults/usr/local/etc/fscd.conf
share/trueos/desktop-defaults/usr/local/etc/init.d/trueosinit
share/trueos/desktop-defaults/usr/local/etc/pkg.conf
share/trueos/desktop-defaults/usr/local/etc/polkit-1/localauthority/50-local.d/system.pkla
share/trueos/desktop-defaults/usr/local/etc/rc.d/trueosinit
share/trueos/desktop-defaults/usr/local/etc/screenrc
share/trueos/desktop-defaults/usr/local/etc/sudoers
share/trueos/desktop-defaults/usr/local/etc/trueos.conf
share/trueos/desktop-defaults/usr/local/etc/volman.conf
share/trueos/desktop-defaults/usr/local/etc/xdg/autostart/pc-mixer.desktop
share/trueos/desktop-defaults/usr/local/etc/xdg/autostart/pc-mounttray.desktop
share/trueos/desktop-defaults/usr/local/etc/xdg/autostart/pc-welcome.desktop
share/trueos/desktop-defaults/usr/local/share/trueos/fluxboxinit/init
share/trueos/desktop-defaults/usr/local/share/trueos/fluxboxinit/keys
share/trueos/desktop-defaults/usr/local/share/trueos/fluxboxinit/menu
share/trueos/desktop-defaults/usr/local/share/trueos/fluxboxinit/overlay
share/trueos/desktop-defaults/usr/local/share/trueos/fluxboxinit/startup
share/trueos/desktop-defaults/usr/local/share/trueos/icons/gvim.png
share/trueos/desktop-defaults/usr/local/share/trueos/icons/hwinfo.png
share/trueos/desktop-defaults/usr/local/share/trueos/icons/input-keyboard.png
share/trueos/desktop-defaults/usr/local/share/trueos/icons/ldap.png
share/trueos/desktop-defaults/usr/local/share/trueos/icons/mythtv.png
share/trueos/desktop-defaults/usr/local/share/trueos/icons/preferences-desktop-theme.png
share/trueos/desktop-defaults/usr/local/share/trueos/icons/preferences-desktop.png
share/trueos/desktop-defaults/usr/local/share/trueos/icons/soundconfig.png
share/trueos/desktop-defaults/usr/local/share/trueos/icons/system-software-update.png
share/trueos/desktop-defaults/usr/local/share/trueos/icons/system-users.png
share/trueos/desktop-defaults/usr/local/share/trueos/icons/utilities-system-monitor.png
share/trueos/desktop-defaults/usr/local/share/trueos/scripts/first-boot.sh
share/trueos/desktop-defaults/usr/local/share/trueos/scripts/lidchanged.sh
share/trueos/desktop-defaults/usr/local/share/trueos/scripts/move-vbox.sh
share/trueos/desktop-defaults/usr/local/share/trueos/scripts/mythtv.sh
share/trueos/desktop-defaults/usr/local/share/trueos/scripts/openhandbook.sh
share/trueos/desktop-defaults/usr/local/share/trueos/scripts/portsnap-noterm.sh
share/trueos/desktop-defaults/usr/local/share/trueos/scripts/set-lang.sh
share/trueos/desktop-defaults/usr/local/share/trueos/scripts/sound-dev-changed.sh
share/trueos/desktop-defaults/usr/local/share/trueos/scripts/vbox-check
share/trueos/desktop-defaults/usr/local/share/trueos/scripts/waitforsystray.sh
share/trueos/desktop-defaults/usr/local/share/trueos/xdg-auto/checkwelcome.sh
share/trueos/desktop-defaults/usr/local/share/trueos/xdg-auto/life-preserver.sh
share/trueos/desktop-defaults/usr/local/share/trueos/xdg-auto/pc-mixer.sh
share/trueos/desktop-defaults/usr/local/share/trueos/xdg-auto/pc-mounttray.sh
share/trueos/desktop-defaults/usr/local/share/trueos/xdg-auto/pc-trayapps.sh
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/a_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/a_1600x900.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/b_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/b_1600x900.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/c_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/c_1600x900.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/d_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/d_1600x900.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/e_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/e_1600x900.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/f_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/f_1600x900.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/g_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/FreeBSD/g_1600x900.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1024x600_v2.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1024x768_v2.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1280x1024_v2.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1280x800_v2.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1440x900_v2.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1600x1200_v2.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1920x1080_v2.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/PCBSD_10_Joule_WP_1920x1200_v2.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/PCBSD_10_Joule_WP_600x400_v2.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/PCBSD_10_Joule_WP_800x600_v2.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/PCBSD_Enterprise_WP_1920x1080.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/a2_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/a_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/b_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/c_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/d_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/e_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/f2_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/f_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/PCBSD/g_1280x1024.jpg
share/trueos/desktop-defaults/usr/local/share/wallpapers/TrueOS/pcdm-trueos-1080.png
share/trueos/desktop-defaults/usr/local/share/wallpapers/TrueOS/pcdm-trueos-4k.png
share/trueos/desktop-defaults/usr/local/share/wallpapers/TrueOS/trueos-1-1080.png
share/trueos/desktop-defaults/usr/local/share/wallpapers/TrueOS/trueos-1-4k.png
share/trueos/desktop-defaults/usr/local/share/wallpapers/TrueOS/trueos-2-1080.png
share/trueos/desktop-defaults/usr/local/share/wallpapers/TrueOS/trueos-2-4k.png
share/trueos/desktop-defaults/usr/local/share/wallpapers/TrueOS/trueos-current-underconstruction-1080.png
share/trueos/desktop-defaults/usr/local/share/wallpapers/TrueOS/trueos-current-underconstruction-4k.png
share/trueos/desktop-defaults/usr/local/share/wallpapers/trueos-default.jpg
share/trueos/desktop-defaults/usr/share/skel/.bashrc
share/trueos/desktop-defaults/usr/share/skel/.cache/dconf/user
share/trueos/desktop-defaults/usr/share/skel/.config/autostart/gnome-fallback-mount-helper.desktop
share/trueos/desktop-defaults/usr/share/skel/.config/autostart/gnome-settings-daemon.desktop
share/trueos/desktop-defaults/usr/share/skel/.config/autostart/ibus.desktop
share/trueos/desktop-defaults/usr/share/skel/.config/dconf/user
share/trueos/desktop-defaults/usr/share/skel/.config/qtcurve.gtk-icons
share/trueos/desktop-defaults/usr/share/skel/.config/qtcurve/gtk-icons
share/trueos/desktop-defaults/usr/share/skel/.config/qtcurve/stylerc
share/trueos/desktop-defaults/usr/share/skel/.config/user-dirs.dirs
share/trueos/desktop-defaults/usr/share/skel/.fluxbox/menu
share/trueos/desktop-defaults/usr/share/skel/.fluxbox/overlay
share/trueos/desktop-defaults/usr/share/skel/.fluxbox/startup
share/trueos/desktop-defaults/usr/share/skel/.icons/default/index.theme
share/trueos/desktop-defaults/usr/share/skel/.screenrc
share/trueos/desktop-defaults/usr/share/skel/.vimrc
share/trueos/desktop-defaults/usr/share/skel/.xprofile
share/trueos/desktop-defaults/usr/share/skel/.xscreensaver
share/trueos/desktop-defaults/usr/share/skel/.xsession
share/trueos/desktop-defaults/usr/share/skel/Desktop/.directory
share/trueos/desktop-defaults/usr/share/skel/Documents/.directory
share/trueos/desktop-defaults/usr/share/skel/Downloads/.directory
share/trueos/desktop-defaults/usr/share/skel/Music/.directory
share/trueos/desktop-defaults/usr/share/skel/Pictures/.keep_me
share/trueos/desktop-defaults/usr/share/skel/Videos/.directory
share/trueos/desktop-defaults/usr/share/skel/dot.Xresources
share/trueos/desktop-defaults/usr/share/skel/dot.cshrc
share/trueos/desktop-defaults/usr/share/skel/dot.face.icon
share/trueos/desktop-defaults/usr/share/skel/dot.ideskrc
share/trueos/desktop-defaults/usr/share/skel/dot.login
share/trueos/desktop-defaults/usr/share/skel/dot.loginIcon.png
share/trueos/desktop-defaults/usr/share/skel/dot.login_conf
share/trueos/desktop-defaults/usr/share/skel/dot.mail_aliases
share/trueos/desktop-defaults/usr/share/skel/dot.mailrc
share/trueos/desktop-defaults/usr/share/skel/dot.profile
share/trueos/desktop-defaults/usr/share/skel/dot.rhosts
share/trueos/desktop-defaults/usr/share/skel/dot.shrc
share/trueos/desktop-defaults/usr/share/skel/dot.zshrc
share/trueos/scripts/mkzfsdir.sh
share/trueos/scripts/reset-firewall
share/trueos/scripts/reset-firewall-tor
share/trueos/scripts/rmzfsdir.sh
share/trueos/scripts/sys-init.sh
share/trueos/server-defaults/boot/grub/themes/trueos/README
share/trueos/server-defaults/boot/grub/themes/trueos/background.png
share/trueos/server-defaults/boot/grub/themes/trueos/blob_w.png
share/trueos/server-defaults/boot/grub/themes/trueos/boot_menu_c.png
share/trueos/server-defaults/boot/grub/themes/trueos/boot_menu_e.png
share/trueos/server-defaults/boot/grub/themes/trueos/boot_menu_n.png
share/trueos/server-defaults/boot/grub/themes/trueos/boot_menu_ne.png
share/trueos/server-defaults/boot/grub/themes/trueos/boot_menu_nw.png
share/trueos/server-defaults/boot/grub/themes/trueos/boot_menu_s.png
share/trueos/server-defaults/boot/grub/themes/trueos/boot_menu_se.png
share/trueos/server-defaults/boot/grub/themes/trueos/boot_menu_sw.png
share/trueos/server-defaults/boot/grub/themes/trueos/boot_menu_w.png
share/trueos/server-defaults/boot/grub/themes/trueos/dejavu_10.pf2
share/trueos/server-defaults/boot/grub/themes/trueos/dejavu_12.pf2
share/trueos/server-defaults/boot/grub/themes/trueos/dejavu_14.pf2
share/trueos/server-defaults/boot/grub/themes/trueos/dejavu_16.pf2
share/trueos/server-defaults/boot/grub/themes/trueos/dejavu_bold_14.pf2
share/trueos/server-defaults/boot/grub/themes/trueos/slider_c.png
share/trueos/server-defaults/boot/grub/themes/trueos/slider_n.png
share/trueos/server-defaults/boot/grub/themes/trueos/slider_s.png
share/trueos/server-defaults/boot/grub/themes/trueos/terminal_box_c.png
share/trueos/server-defaults/boot/grub/themes/trueos/terminal_box_e.png
share/trueos/server-defaults/boot/grub/themes/trueos/terminal_box_n.png
share/trueos/server-defaults/boot/grub/themes/trueos/terminal_box_ne.png
share/trueos/server-defaults/boot/grub/themes/trueos/terminal_box_nw.png
share/trueos/server-defaults/boot/grub/themes/trueos/terminal_box_s.png
share/trueos/server-defaults/boot/grub/themes/trueos/terminal_box_se.png
share/trueos/server-defaults/boot/grub/themes/trueos/terminal_box_sw.png
share/trueos/server-defaults/boot/grub/themes/trueos/terminal_box_w.png
share/trueos/server-defaults/boot/grub/themes/trueos/theme.txt
share/trueos/server-defaults/boot/grub/trueosfont.pf2
share/trueos/server-defaults/root/.bashrc
share/trueos/server-defaults/root/.cshrc
share/trueos/server-defaults/root/.profile
share/trueos/server-defaults/root/.vimrc
share/trueos/server-defaults/usr/local/etc/beadm.conf
share/trueos/server-defaults/usr/local/etc/default/grub
share/trueos/server-defaults/usr/local/etc/pkg.conf
share/trueos/server-defaults/usr/local/etc/screenrc
share/trueos/server-defaults/usr/local/etc/trueos.conf
share/trueos/server-defaults/usr/share/skel/.bashrc
share/trueos/server-defaults/usr/share/skel/.vimrc
share/trueos/server-defaults/usr/share/skel/dot.cshrc
share/trueos/server-defaults/usr/share/skel/dot.profile
share/trueos/server-defaults/usr/share/skel/dot.shrc
share/trueos/xstartup/disable-beep.sh
share/trueos/xstartup/enable-inputmethod.sh
share/trueos/xstartup/enable-xprofile.sh
share/trueos/xstartup/gpg-agent.sh
share/trueos/xstartup/ssh-agent.sh
