# security/vpnc/pkg-plist
sbin/vpnc
etc/vpnc.conf.sample
etc/rc.d/vpnc.sh.sample
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%@dirrm %%DOCSDIR%%
