# www/xitami/pkg-plist
etc/rc.d/xitami.sh.sample
share/doc/xitami/html/index.htm
share/doc/xitami/html/index1.htm
share/doc/xitami/html/index10.htm
share/doc/xitami/html/index11.htm
share/doc/xitami/html/index12.htm
share/doc/xitami/html/index13.htm
share/doc/xitami/html/index14.htm
share/doc/xitami/html/index15.htm
share/doc/xitami/html/index16.htm
share/doc/xitami/html/index2.htm
share/doc/xitami/html/index3.htm
share/doc/xitami/html/index4.htm
share/doc/xitami/html/index5.htm
share/doc/xitami/html/index6.htm
share/doc/xitami/html/index7.htm
share/doc/xitami/html/index8.htm
share/doc/xitami/html/index9.htm
share/doc/xitami/html/license.htm
xitami/cgi-bin/testcgi
xitami/defaults.cfg.sample
xitami/webpages/default.html.sample
xitami/xitami
xitami/xitami.aut
xitami/xitami.cfg
@dirrm share/doc/xitami/html
@dirrm share/doc/xitami
@unexec rmdir %D/xitami/webpages 2>/dev/null || true
@unexec rmdir %D/xitami/cgi-bin 2>/dev/null || true
@unexec rmdir %D/xitami 2>/dev/null || true
