# irc/ircd-rusnet/pkg-plist
bin/irc
etc/ircd-rusnet/codepages/ALT2UNI
etc/ircd-rusnet/codepages/ISO2UNI
etc/ircd-rusnet/codepages/KOI2UNI
etc/ircd-rusnet/codepages/MAC2UNI
etc/ircd-rusnet/codepages/MRO2UNI
etc/ircd-rusnet/codepages/TRN2UNI
etc/ircd-rusnet/codepages/UNI2ALT
etc/ircd-rusnet/codepages/UNI2ISO
etc/ircd-rusnet/codepages/UNI2KOI
etc/ircd-rusnet/codepages/UNI2MAC
etc/ircd-rusnet/codepages/UNI2MRO
etc/ircd-rusnet/codepages/UNI2TRN
etc/ircd-rusnet/codepages/UNI2WIN
etc/ircd-rusnet/codepages/WIN2UNI
etc/ircd-rusnet/example.conf
etc/ircd-rusnet/example.rusnet.conf
etc/ircd-rusnet/example.tkserv.access
etc/ircd-rusnet/iauth.conf
etc/ircd-rusnet/kills.conf
etc/ircd-rusnet/ircd.m4
etc/ircd-rusnet/ircd.motd.sample
etc/rc.d/ircd.sh.sample
sbin/ircd
sbin/iauth
sbin/chkconf
sbin/ircd-mkpasswd
sbin/ircdwatch
sbin/tkserv
%%PORTDOCS%%%%DOCSDIR%%/CHANGED
%%PORTDOCS%%%%DOCSDIR%%/FAQ
%%PORTDOCS%%%%DOCSDIR%%/Authors
%%PORTDOCS%%%%DOCSDIR%%/BUGS
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/Etiquette
%%PORTDOCS%%%%DOCSDIR%%/INSTALL.appendix
%%PORTDOCS%%%%DOCSDIR%%/INSTALL.txt
%%PORTDOCS%%%%DOCSDIR%%/LICENSE
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/RELEASE_NOTES
%%PORTDOCS%%%%DOCSDIR%%/SERVICE.txt
%%PORTDOCS%%%%DOCSDIR%%/alt-irc-faq
%%PORTDOCS%%%%DOCSDIR%%/iauth-internals.txt
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm etc/ircd-rusnet/codepages
@dirrm etc/ircd-rusnet
@unexec rmdir /var/log/ircd 2>/dev/null || true
@unexec rmdir /var/run/ircd 2>/dev/null || true
