# mail/tpop3d/files/tpop3d.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: tpop3d
# REQUIRE:  postgresql mysql LOGIN
#
# Add the following line to /etc/rc.conf to enable tpop3d:
#tpop3d_enable="YES"
#
# You may also wish to use the following variables to fine-tune startup:
#tpop3d_flags="..."

. /etc/rc.subr

name="tpop3d"
rcvar=tpop3d_enable


stop_postcmd=stop_postcmd

stop_postcmd()
{
  rm -f $pidfile
}

# set defaults


command=%%PREFIX%%/sbin/tpop3d
required_files=%%PREFIX%%/etc/tpop3d.conf
pidfile=/var/run/tpop3d.pid

load_rc_config $name

: ${tpop3d_enable="NO"}
: ${tpop3d_flags="-f ${required_files} -p ${pidfile}"}

run_rc_command "$1"
