# net/dbeacon/files/openrc-dbeacon.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf to enable dbeacon.
#
# dbeacon_flags=""	Extra flags passed to start command.
#
# See dbeacon(1) for flags.

dbeacon_flags=${dbeacon_flags-""}

name="dbeacon"
command="%%PREFIX%%/bin/${name}"

depend() {
	keyword -shutdown
}
