# net/ptpd2-devel/files/ptpd2.in
#!/bin/sh
#
# $FreeBSD: /tmp/pcvs/ports/net/ptpd2-devel/files/ptpd2.in,v 1.1 2012-04-11 21:00:51 gnn Exp $
#
# PROVIDE: ptpd
# REQUIRE: NETWORKING DAEMON
# KEYWORD: nojail
#
# Add the following lines to /etc/rc.conf to enable ptpd
#
# ptpd_enable (bool):	Set to "NO" by default
#			Set it to "YES" to enable ptpd

. /etc/rc.subr

name="ptpd2"
rcvar=`set_rcvar`

command="%%PREFIX%%/sbin/${name}"

ptpd2_enable=${ptpd2_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"
