# audio/autocd/pkg-plist
bin/autocd
bin/cdctl
bin/pause
bin/eject
etc/rc.d/autocd.sh
