# audio/autocd/pkg-plist
bin/acdplay
bin/autocd
bin/cdctl
bin/pause
bin/eject
etc/rc.d/autocd.sh
%%PORTDOCS%%%%DOCSDIR%%/README.rus.html
%%PORTDOCS%%@dirrm %%DOCSDIR%%
