# audio/autocd/pkg-plist
bin/autocd
bin/cdctl
bin/pause
bin/eject
etc/rc.d/autocd.sh
%%PORTDOCS%%share/doc/autocd/README.rus.html
@dirrm share/doc/autocd
