# mail/mailutils/files/comsatd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: comsatd
# REQUIRE: DAEMON
# BEFORE: mail
# KEYWORD: shutdown
#

. /etc/rc.subr

name="comsatd"
rcvar=comsatd_enable
command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/${name}.pid"

load_rc_config ${name}

: ${comsatd_enable="NO"}

run_rc_command "$1"
