# mail/mailutils/files/openrc-imap4d.in
#!/sbin/openrc-run

name="imap4d"
command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/${name}.pid"

depend() {
	keyword -shutdown
}
