# mail/mailutils/files/openrc-pop3d.in
#!/sbin/openrc-run

name="pop3d"
command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/${name}.pid"

depend() {
	keyword -shutdown
}
