# mail/mailutils/files/openrc-comstad.in
#!/sbin/openrc-run

name="comsatd"
command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/${name}.pid"

depend() {
	keyword -shutdown
}
