# mail/mailutils/files/pop3d.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: pop3d
# REQUIRE: DAEMON
# BEFORE: mail
# KEYWORD: shutdown
#

. /etc/rc.subr

name="pop3d"
rcvar=pop3d_enable

: ${pop3d_enable="NO"}

command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/${name}.pid"

required_files="%%PREFIX%%/etc/mailutils.conf %%PREFIX%%/etc/mailutils/${name}"

extra_commands="configlint"
configlint_cmd="${name} --config-lint"

load_rc_config ${name}
run_rc_command "$1"
