# sysutils/apcupsd/pkg-plist
sbin/apcupsd
sbin/apcaccess
sbin/apcnisd
sbin/apctest
sbin/smtp
%%POWERFL%%sbin/powerflute
etc/apcupsd/apcupsd.conf.net-master.sample
etc/apcupsd/apcupsd.conf.net-slave.sample
@unexec if cmp -s %D/etc/apcupsd/apcupsd.conf.sample %D/etc/apcupsd/apcupsd.conf; then rm -f %D/etc/apcupsd/apcupsd.conf; fi
etc/apcupsd/apcupsd.conf.sample
etc/apcupsd/apccontrol
etc/apcupsd/changeme
etc/apcupsd/commfailure
etc/apcupsd/commok
etc/apcupsd/mainsback
etc/apcupsd/masterconnect
etc/apcupsd/mastertimeout
etc/apcupsd/onbattery
%%CGI%%@unexec if cmp -s %D/etc/apcupsd/apcupsd.css.sample %D/etc/apcupsd/apcupsd.css; then rm -f %D/etc/apcupsd/apcupsd.css; fi
%%CGI%%etc/apcupsd/apcupsd.css.sample
%%CGI%%@unexec if cmp -s %D/etc/apcupsd/host.conf.sample %D/etc/apcupsd/host.conf; then rm -f %D/etc/apcupsd/host.conf; fi
%%CGI%%etc/apcupsd/hosts.conf.sample
%%CGI%%@unexec if cmp -s %D/etc/apcupsd/multimon.conf.sample %D/etc/apcupsd/multimon.conf; then rm -f %D/etc/apcupsd/multimon.conf; fi
%%CGI%%etc/apcupsd/multimon.conf.sample
%%CGI%%etc/apcupsd/cgi/multimon.cgi
%%CGI%%etc/apcupsd/cgi/upsfstats.cgi
%%CGI%%etc/apcupsd/cgi/upsimage.cgi
%%CGI%%etc/apcupsd/cgi/upsstats.cgi
etc/rc.d/%%STARTUP_SCRIPT%%
%%CGI%%@dirrm etc/apcupsd/cgi
@unexec rmdir %D/etc/apcupsd 2>/dev/null || /usr/bin/true
