# sysutils/apcupsd/pkg-plist
sbin/apcupsd
sbin/apcaccess
sbin/apcnisd
sbin/apctest
%%POWERFL%%sbin/powerflute
etc/apcupsd/apccontrol
@unexec if [ -f %D/etc/apcupsd/apcupsd.conf ]; then cmp -s %D/etc/apcupsd/apcupsd.conf.sample %D/etc/apcupsd/apcupsd.conf && rm -f %D/etc/apcupsd/apcupsd.conf || echo "If you are permanently removing this port, you should do a ``rm ${PKG_PREFIX}/etc/apcupsd/apcupsd.conf`` to remove config files left." | fmt ; fi
etc/apcupsd/apcupsd.conf.sample
etc/apcupsd/changeme
etc/apcupsd/commfailure
etc/apcupsd/commok
etc/apcupsd/onbattery
etc/apcupsd/mainsback
@unexec if [ -f %D/etc/apcupsd/hosts.conf ]; then cmp -s %D/etc/apcupsd/hosts.conf.sample %D/etc/apcupsd/hosts.conf && rm -f %D/etc/apcupsd/hosts.conf || echo "If you are permanently removing this port, you should do a ``rm ${PKG_PREFIX}/etc/apcupsd/hosts.conf`` to remove config files left." | fmt ; fi
etc/apcupsd/hosts.conf.sample
@unexec if [ -f %D/etc/apcupsd/multimon.conf ]; then cmp -s %D/etc/apcupsd/multimon.conf.sample %D/etc/apcupsd/multimon.conf && rm -f %D/etc/apcupsd/multimon.conf || echo "If you are permanently removing this port, you should do a ``rm ${PKG_PREFIX}/etc/apcupsd/multimon.conf`` to remove config files left." | fmt ; fi
etc/apcupsd/multimon.conf.sample
etc/apcupsd/cgi/multimon.cgi
etc/apcupsd/cgi/multimoncss.cgi
etc/apcupsd/cgi/upsfstats.cgi
etc/apcupsd/cgi/upsimage.cgi
etc/apcupsd/cgi/upsstats.cgi
etc/rc.d/%%STARTUP_SCRIPT%%
share/locale/it/LC_MESSAGES/apcupsd.mo
%%PORTDOCS%%share/doc/apcupsd/apcupsd.pdf
%%PORTDOCS%%share/doc/apcupsd/apcaccess.html
%%PORTDOCS%%share/doc/apcupsd/apcnisd.html
%%PORTDOCS%%share/doc/apcupsd/apctest.html
%%PORTDOCS%%share/doc/apcupsd/apcupsd-styles.css
%%PORTDOCS%%share/doc/apcupsd/apcupsd.apf
%%PORTDOCS%%share/doc/apcupsd/apcupsd.book
%%PORTDOCS%%share/doc/apcupsd/apcupsd.gif
%%PORTDOCS%%share/doc/apcupsd/apcupsd.html
%%PORTDOCS%%share/doc/apcupsd/batteries.html
%%PORTDOCS%%share/doc/apcupsd/bugs.html
%%PORTDOCS%%share/doc/apcupsd/cables.html
%%PORTDOCS%%share/doc/apcupsd/cgiprogs.html
%%PORTDOCS%%share/doc/apcupsd/charging.gif
%%PORTDOCS%%share/doc/apcupsd/config-examples.html
%%PORTDOCS%%share/doc/apcupsd/configure.html
%%PORTDOCS%%share/doc/apcupsd/data.html
%%PORTDOCS%%share/doc/apcupsd/eprom.html
%%PORTDOCS%%share/doc/apcupsd/events.html
%%PORTDOCS%%share/doc/apcupsd/faq.html
%%PORTDOCS%%share/doc/apcupsd/index.html
%%PORTDOCS%%share/doc/apcupsd/install.html
%%PORTDOCS%%share/doc/apcupsd/invoking.html
%%PORTDOCS%%share/doc/apcupsd/kernel_config.html
%%PORTDOCS%%share/doc/apcupsd/license.html
%%PORTDOCS%%share/doc/apcupsd/logging.html
%%PORTDOCS%%share/doc/apcupsd/master-slave.html
%%PORTDOCS%%share/doc/apcupsd/multiUPS.html
%%PORTDOCS%%share/doc/apcupsd/multimon.gif
%%PORTDOCS%%share/doc/apcupsd/new_features_3_7_0.html
%%PORTDOCS%%share/doc/apcupsd/new_features_3_8_0.html
%%PORTDOCS%%share/doc/apcupsd/new_features_3_8_1.html
%%PORTDOCS%%share/doc/apcupsd/new_features_3_8_2.html
%%PORTDOCS%%share/doc/apcupsd/new_features_3_8_3.html
%%PORTDOCS%%share/doc/apcupsd/new_features_3_8_4.html
%%PORTDOCS%%share/doc/apcupsd/new_features_3_8_5.html
%%PORTDOCS%%share/doc/apcupsd/oldversions.html
%%PORTDOCS%%share/doc/apcupsd/onbatt.gif
%%PORTDOCS%%share/doc/apcupsd/online.gif
%%PORTDOCS%%share/doc/apcupsd/security.html
%%PORTDOCS%%share/doc/apcupsd/shutdown.html
%%PORTDOCS%%share/doc/apcupsd/status.gif
%%PORTDOCS%%share/doc/apcupsd/status.html
%%PORTDOCS%%share/doc/apcupsd/stopping.html
%%PORTDOCS%%share/doc/apcupsd/testing.html
%%PORTDOCS%%share/doc/apcupsd/thanks.gif
%%PORTDOCS%%share/doc/apcupsd/thanks.html
%%PORTDOCS%%share/doc/apcupsd/troubles.html
%%PORTDOCS%%share/doc/apcupsd/upgrading_to_3_7_0.html
%%PORTDOCS%%share/doc/apcupsd/upsbible.html
%%PORTDOCS%%share/doc/apcupsd/usb.html
%%PORTDOCS%%share/doc/apcupsd/win32.html
%%PORTDOCS%%share/doc/apcupsd/wininstall1.gif
%%PORTDOCS%%share/doc/apcupsd/wininstall2.gif
%%PORTDOCS%%share/doc/apcupsd/wininstall3.gif
%%PORTDOCS%%share/doc/apcupsd/wininstall4.gif
%%PORTDOCS%%share/doc/apcupsd/wininstall5.gif
%%PORTDOCS%%share/doc/apcupsd/wininstall6.gif
%%PORTDOCS%%share/doc/apcupsd/wininstall7.gif
@dirrm etc/apcupsd/cgi
@dirrm etc/apcupsd
%%PORTDOCS%%@dirrm share/doc/apcupsd
