# www/webfs/files/openrc-webfsd.in
#!/sbin/openrc-run

name="webfsd"

command=%%PREFIX%%/sbin/webfsd

depend(){
	keyword -shutdown
}
