# mail/mailman/files/mailman.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: mailman
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable or configure mailman:
#
# mailman_enable (bool):        Set it to "YES" to enable mailman.
#                               Default is "NO".
#
# mailman_fix_perms (bool):     Set it to "NO" to disable running
#                               check_perms on start or post-install.
#                               Default is "YES", but this is unsuitable
#                               for slow machines with many lists.

. /etc/rc.subr

name="mailman"
rcvar="mailman_enable"
required_dirs="%%MAILMANDIR%%/lists/mailman"
command_interpreter="%%PYTHON_CMD%%"
command="%%MAILMANDIR%%/bin/mailmanctl"
command_args="start"
pidfile="%%MAILMANDIR%%/data/master-qrunner.pid"

load_rc_config $name

: ${mailman_enable:="NO"}
: ${mailman_fix_perms:="YES"}
: ${mailman_flags:="-s -q"}

start_precmd=${name}_prestart
extra_commands="reload status"

mailman_prestart() {
	if checkyesno mailman_fix_perms ; then 
	    check_startmsgs && echo "Fixing ${name} permissions:"
	    env LC_ALL=C %%MAILMANDIR%%/bin/check_perms -f ${rc_debug:+-v}
	fi
}

run_rc_command "$1"
