# irc/ircd-hybrid/pkg-plist
bin/ircd
bin/mkpasswd
bin/viconf
bin/viklines
bin/vimotd
etc/ircd/ircd.conf.sample
etc/ircd/ircd.motd.sample
etc/ircd/kline.conf.sample
etc/ircd/opers.txt.sample
etc/rc.d/ircd.sh
%%PORTDOCS%%/doc/ircd-hybrid/CIDR.txt
%%PORTDOCS%%/doc/ircd-hybrid/INDEX
%%PORTDOCS%%/doc/ircd-hybrid/Makefile
%%PORTDOCS%%/doc/ircd-hybrid/README.TSora
%%PORTDOCS%%/doc/ircd-hybrid/Tao-of-IRC.940110
%%PORTDOCS%%/doc/ircd-hybrid/blalloc.txt
%%PORTDOCS%%/doc/ircd-hybrid/example.conf
%%PORTDOCS%%/doc/ircd-hybrid/example.conf.trillian
%%PORTDOCS%%/doc/ircd-hybrid/ircd.8
%%PORTDOCS%%/doc/ircd-hybrid/mtrie.txt
%%PORTDOCS%%/doc/ircd-hybrid/old/Authors
%%PORTDOCS%%/doc/ircd-hybrid/old/Etiquette
%%PORTDOCS%%/doc/ircd-hybrid/old/README
%%PORTDOCS%%/doc/ircd-hybrid/old/US-Admin/Networking
%%PORTDOCS%%/doc/ircd-hybrid/old/US-Admin/Operators
%%PORTDOCS%%/doc/ircd-hybrid/old/US-Admin/README
%%PORTDOCS%%/doc/ircd-hybrid/operguide.txt
%%PORTDOCS%%/doc/ircd-hybrid/opermyth.txt
%%PORTDOCS%%/doc/ircd-hybrid/rfc1459.txt
%%PORTDOCS%%/doc/ircd-hybrid/server-version-info
%%PORTDOCS%%/doc/ircd-hybrid/simple.conf
%%PORTDOCS%%/doc/ircd-hybrid/ts5.txt
%%PORTDOCS%%/doc/ircd-hybrid/yline.txt
%%PORTDOCS%%@dirrm share/doc/ircd-hybrid/old/US-Admin
%%PORTDOCS%%@dirrm share/doc/ircd-hybrid/old
%%PORTDOCS%%@dirrm share/doc/ircd-hybrid
@unexec rmdir %D/etc/ircd 2>/dev/null || echo "If you are permanently removing this port, you should do a ``rm -rf ${PKG_PREFIX}/etc/ircd`` to remove any configuration files left." | fmt
