# irc/ircd-hybrid/pkg-plist
bin/convertconf
bin/convertilines
bin/convertklines
bin/encspeed
bin/ircd
bin/mkpasswd
bin/servlink
bin/viconf
bin/viklines
bin/vimotd
etc/ircd-hybrid/example.conf
etc/ircd-hybrid/example.efnet.conf
etc/ircd-hybrid/ircd.motd.sample
etc/ircd-hybrid/kline.conf.sample
etc/ircd-hybrid/dline.conf.sample
etc/ircd-hybrid/.convertconf-example.conf
etc/rc.d/ircd-hybrid.sh
include/ircd-hybrid-7
share/ircd-hybrid-%%PORTVERSION%%/help/opers/accept
share/ircd-hybrid-%%PORTVERSION%%/help/opers/admin
share/ircd-hybrid-%%PORTVERSION%%/help/opers/away
share/ircd-hybrid-%%PORTVERSION%%/help/opers/capab
share/ircd-hybrid-%%PORTVERSION%%/help/opers/cburst
share/ircd-hybrid-%%PORTVERSION%%/help/opers/challenge
share/ircd-hybrid-%%PORTVERSION%%/help/opers/cjoin
share/ircd-hybrid-%%PORTVERSION%%/help/opers/client
share/ircd-hybrid-%%PORTVERSION%%/help/opers/close
share/ircd-hybrid-%%PORTVERSION%%/help/opers/cmode
share/ircd-hybrid-%%PORTVERSION%%/help/opers/connect
share/ircd-hybrid-%%PORTVERSION%%/help/opers/credits
share/ircd-hybrid-%%PORTVERSION%%/help/opers/cryptlink
share/ircd-hybrid-%%PORTVERSION%%/help/opers/die
share/ircd-hybrid-%%PORTVERSION%%/help/opers/dline
share/ircd-hybrid-%%PORTVERSION%%/help/opers/dmem
share/ircd-hybrid-%%PORTVERSION%%/help/opers/drop
share/ircd-hybrid-%%PORTVERSION%%/help/opers/eob
share/ircd-hybrid-%%PORTVERSION%%/help/opers/error
share/ircd-hybrid-%%PORTVERSION%%/help/opers/gline
share/ircd-hybrid-%%PORTVERSION%%/help/opers/help
share/ircd-hybrid-%%PORTVERSION%%/help/opers/index
share/ircd-hybrid-%%PORTVERSION%%/help/opers/info
share/ircd-hybrid-%%PORTVERSION%%/help/opers/invite
share/ircd-hybrid-%%PORTVERSION%%/help/opers/ison
share/ircd-hybrid-%%PORTVERSION%%/help/opers/join
share/ircd-hybrid-%%PORTVERSION%%/help/opers/kick
share/ircd-hybrid-%%PORTVERSION%%/help/opers/kill
share/ircd-hybrid-%%PORTVERSION%%/help/opers/kline
share/ircd-hybrid-%%PORTVERSION%%/help/opers/knock
share/ircd-hybrid-%%PORTVERSION%%/help/opers/knockll
share/ircd-hybrid-%%PORTVERSION%%/help/opers/links
share/ircd-hybrid-%%PORTVERSION%%/help/opers/list
share/ircd-hybrid-%%PORTVERSION%%/help/opers/lljoin
share/ircd-hybrid-%%PORTVERSION%%/help/opers/llnick
share/ircd-hybrid-%%PORTVERSION%%/help/opers/locops
share/ircd-hybrid-%%PORTVERSION%%/help/opers/lusers
share/ircd-hybrid-%%PORTVERSION%%/help/opers/modlist
share/ircd-hybrid-%%PORTVERSION%%/help/opers/modload
share/ircd-hybrid-%%PORTVERSION%%/help/opers/modrestart
share/ircd-hybrid-%%PORTVERSION%%/help/opers/modunload
share/ircd-hybrid-%%PORTVERSION%%/help/opers/motd
share/ircd-hybrid-%%PORTVERSION%%/help/opers/names
share/ircd-hybrid-%%PORTVERSION%%/help/opers/nburst
share/ircd-hybrid-%%PORTVERSION%%/help/opers/nick
share/ircd-hybrid-%%PORTVERSION%%/help/opers/notice
share/ircd-hybrid-%%PORTVERSION%%/help/opers/oper
share/ircd-hybrid-%%PORTVERSION%%/help/opers/operwall
share/ircd-hybrid-%%PORTVERSION%%/help/opers/part
share/ircd-hybrid-%%PORTVERSION%%/help/opers/pass
share/ircd-hybrid-%%PORTVERSION%%/help/opers/ping
share/ircd-hybrid-%%PORTVERSION%%/help/opers/pong
share/ircd-hybrid-%%PORTVERSION%%/help/opers/post
share/ircd-hybrid-%%PORTVERSION%%/help/opers/privmsg
share/ircd-hybrid-%%PORTVERSION%%/help/opers/quit
share/ircd-hybrid-%%PORTVERSION%%/help/opers/rehash
share/ircd-hybrid-%%PORTVERSION%%/help/opers/restart
share/ircd-hybrid-%%PORTVERSION%%/help/opers/resv
share/ircd-hybrid-%%PORTVERSION%%/help/opers/server
share/ircd-hybrid-%%PORTVERSION%%/help/opers/set
share/ircd-hybrid-%%PORTVERSION%%/help/opers/sjoin
share/ircd-hybrid-%%PORTVERSION%%/help/opers/squit
share/ircd-hybrid-%%PORTVERSION%%/help/opers/stats
share/ircd-hybrid-%%PORTVERSION%%/help/opers/svinfo
share/ircd-hybrid-%%PORTVERSION%%/help/opers/testline
share/ircd-hybrid-%%PORTVERSION%%/help/opers/time
share/ircd-hybrid-%%PORTVERSION%%/help/opers/topic
share/ircd-hybrid-%%PORTVERSION%%/help/opers/trace
share/ircd-hybrid-%%PORTVERSION%%/help/opers/uhelp
share/ircd-hybrid-%%PORTVERSION%%/help/opers/umode
share/ircd-hybrid-%%PORTVERSION%%/help/opers/undline
share/ircd-hybrid-%%PORTVERSION%%/help/opers/ungline
share/ircd-hybrid-%%PORTVERSION%%/help/opers/unkline
share/ircd-hybrid-%%PORTVERSION%%/help/opers/unresv
share/ircd-hybrid-%%PORTVERSION%%/help/opers/user
share/ircd-hybrid-%%PORTVERSION%%/help/opers/userhost
share/ircd-hybrid-%%PORTVERSION%%/help/opers/users
share/ircd-hybrid-%%PORTVERSION%%/help/opers/version
share/ircd-hybrid-%%PORTVERSION%%/help/opers/wallops
share/ircd-hybrid-%%PORTVERSION%%/help/opers/who
share/ircd-hybrid-%%PORTVERSION%%/help/opers/whois
share/ircd-hybrid-%%PORTVERSION%%/help/opers/whowas
share/ircd-hybrid-%%PORTVERSION%%/help/users/accept
share/ircd-hybrid-%%PORTVERSION%%/help/users/admin
share/ircd-hybrid-%%PORTVERSION%%/help/users/away
share/ircd-hybrid-%%PORTVERSION%%/help/users/challenge
share/ircd-hybrid-%%PORTVERSION%%/help/users/cjoin
share/ircd-hybrid-%%PORTVERSION%%/help/users/cmode
share/ircd-hybrid-%%PORTVERSION%%/help/users/credits
share/ircd-hybrid-%%PORTVERSION%%/help/users/error
share/ircd-hybrid-%%PORTVERSION%%/help/users/help
share/ircd-hybrid-%%PORTVERSION%%/help/users/index
share/ircd-hybrid-%%PORTVERSION%%/help/users/info
share/ircd-hybrid-%%PORTVERSION%%/help/users/invite
share/ircd-hybrid-%%PORTVERSION%%/help/users/ison
share/ircd-hybrid-%%PORTVERSION%%/help/users/join
share/ircd-hybrid-%%PORTVERSION%%/help/users/kick
share/ircd-hybrid-%%PORTVERSION%%/help/users/knock
share/ircd-hybrid-%%PORTVERSION%%/help/users/links
share/ircd-hybrid-%%PORTVERSION%%/help/users/list
share/ircd-hybrid-%%PORTVERSION%%/help/users/lusers
share/ircd-hybrid-%%PORTVERSION%%/help/users/motd
share/ircd-hybrid-%%PORTVERSION%%/help/users/names
share/ircd-hybrid-%%PORTVERSION%%/help/users/nick
share/ircd-hybrid-%%PORTVERSION%%/help/users/notice
share/ircd-hybrid-%%PORTVERSION%%/help/users/oper
share/ircd-hybrid-%%PORTVERSION%%/help/users/part
share/ircd-hybrid-%%PORTVERSION%%/help/users/pass
share/ircd-hybrid-%%PORTVERSION%%/help/users/ping
share/ircd-hybrid-%%PORTVERSION%%/help/users/pong
share/ircd-hybrid-%%PORTVERSION%%/help/users/privmsg
share/ircd-hybrid-%%PORTVERSION%%/help/users/quit
share/ircd-hybrid-%%PORTVERSION%%/help/users/stats
share/ircd-hybrid-%%PORTVERSION%%/help/users/time
share/ircd-hybrid-%%PORTVERSION%%/help/users/topic
share/ircd-hybrid-%%PORTVERSION%%/help/users/umode
share/ircd-hybrid-%%PORTVERSION%%/help/users/user
share/ircd-hybrid-%%PORTVERSION%%/help/users/userhost
share/ircd-hybrid-%%PORTVERSION%%/help/users/users
share/ircd-hybrid-%%PORTVERSION%%/help/users/version
share/ircd-hybrid-%%PORTVERSION%%/help/users/who
share/ircd-hybrid-%%PORTVERSION%%/help/users/whois
share/ircd-hybrid-%%PORTVERSION%%/help/users/whowas
share/ircd-hybrid-%%PORTVERSION%%/messages/ayb/LC_MESSAGES/ircd-hybrid.mo
share/ircd-hybrid-%%PORTVERSION%%/messages/custom/LC_MESSAGES/ircd-hybrid.mo
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_accept.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_admin.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_away.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_capab.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_cburst.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_challenge.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_cjoin.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_close.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_connect.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_cryptlink.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_dmem.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_drop.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_encap.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_eob.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_gline.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_help.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_info.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_invite.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_ison.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_join.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_kline.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_knock.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_links.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_list.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_lljoin.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_llnick.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_locops.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_lusers.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_motd.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_names.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_nburst.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_oper.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_operwall.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_pass.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_ping.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_pong.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_post.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_rehash.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_restart.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_resv.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_set.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_stats.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_svinfo.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_testline.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_time.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_topic.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_trace.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_unkline.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_user.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_userhost.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_users.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_version.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_wallops.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_who.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_whois.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_whowas.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/autoload/m_xline.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/m_die.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/m_kick.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/m_kill.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/m_message.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/m_mode.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/m_nick.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/m_part.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/m_quit.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/m_server.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/m_sjoin.so
%%SHAREDMODULES%%share/ircd-hybrid-%%PORTVERSION%%/modules/m_squit.so
@unexec rmdir %D/etc/ircd-hybrid 2>/dev/null || echo "If you are permanently removing this port, you should do a ``rm -rf ${PKG_PREFIX}/etc/ircd-hybrid`` to remove any configuration files left." | fmt
@unexec rmdir /var/log/ircd 2>/dev/null || echo "If you are permanently removing this port, you should do a ``rm -rf /var/log/ircd`` to remove any configuration files left." | fmt
@dirrm share/ircd-hybrid-%%PORTVERSION%%/modules/autoload
@dirrm share/ircd-hybrid-%%PORTVERSION%%/modules
@dirrm share/ircd-hybrid-%%PORTVERSION%%/messages/custom/LC_MESSAGES
@dirrm share/ircd-hybrid-%%PORTVERSION%%/messages/custom
@dirrm share/ircd-hybrid-%%PORTVERSION%%/messages/ayb/LC_MESSAGES
@dirrm share/ircd-hybrid-%%PORTVERSION%%/messages/ayb
@dirrm share/ircd-hybrid-%%PORTVERSION%%/messages
@dirrm share/ircd-hybrid-%%PORTVERSION%%/help/users
@dirrm share/ircd-hybrid-%%PORTVERSION%%/help/opers
@dirrm share/ircd-hybrid-%%PORTVERSION%%/help
@dirrm share/ircd-hybrid-%%PORTVERSION%%
