# irc/ircd-hybrid/pkg-plist
sbin/ircd
bin/buildm4
bin/fixklines
bin/ircd-mkpasswd
bin/viconf
etc/rc.d/ircd.sh
etc/ircd/example.conf.trillian
etc/ircd/ircd.conf
etc/ircd/ircd.motd
etc/ircd/kline.conf
etc/ircd/opers.txt
@dirrm etc/ircd
