# security/sfs/pkg-plist
bin/dirsearch
bin/newaid
bin/rex
bin/rpcc
bin/sfsagent
bin/sfskey
bin/ssu
@unexec if [ -f %D/etc/rc.d/sfscd.sh ]; then cmp -s %D/etc/rc.d/sfscd.sh.sample %D/etc/rc.d/sfscd.sh && rm -f %D/etc/rc.d/sfscd.sh || echo "${PKG_PREFIX}/etc/rc.d/sfscd.sh is different than distribution sample, leaving for manual cleanup." | fmt; fi
etc/rc.d/sfscd.sh.sample
@unexec if [ -f %D/etc/rc.d/sfssd.sh ]; then cmp -s %D/etc/rc.d/sfssd.sh.sample %D/etc/rc.d/sfssd.sh && rm -f %D/etc/rc.d/sfssd.sh || echo "${PKG_PREFIX}/etc/rc.d/sfssd.sh is different than distribution sample, leaving for manual cleanup." | fmt; fi
etc/rc.d/sfssd.sh.sample
etc/sfs/sfsrwsd_config.sample
include/sfs
include/sfs-0.6/aclnt.h
include/sfs-0.6/aes.h
include/sfs-0.6/afsnode.h
include/sfs-0.6/agentconn.h
include/sfs-0.6/agentmisc.h
include/sfs-0.6/aiod.h
include/sfs-0.6/aiod_prot.h
include/sfs-0.6/aios.h
include/sfs-0.6/amisc.h
include/sfs-0.6/arc4.h
include/sfs-0.6/arena.h
include/sfs-0.6/arpc.h
include/sfs-0.6/array.h
include/sfs-0.6/asrv.h
include/sfs-0.6/async.h
include/sfs-0.6/autoconf.h
include/sfs-0.6/axprt.h
include/sfs-0.6/axprt_crypt.h
include/sfs-0.6/backoff.h
include/sfs-0.6/bbuddy.h
include/sfs-0.6/bench.h
include/sfs-0.6/bigint.h
include/sfs-0.6/bitvec.h
include/sfs-0.6/blowfish.h
include/sfs-0.6/callback.h
include/sfs-0.6/cbuf.h
include/sfs-0.6/crypt.h
include/sfs-0.6/crypt_prot.h
include/sfs-0.6/crypt_prot.x
include/sfs-0.6/crypthash.h
include/sfs-0.6/dns.h
include/sfs-0.6/dnsparse.h
include/sfs-0.6/err.h
include/sfs-0.6/esign.h
include/sfs-0.6/fdlim.h
include/sfs-0.6/getfh3.h
include/sfs-0.6/hashcash.h
include/sfs-0.6/ihash.h
include/sfs-0.6/init.h
include/sfs-0.6/itree.h
include/sfs-0.6/keyfunc.h
include/sfs-0.6/list.h
include/sfs-0.6/modalg.h
include/sfs-0.6/mount_prot.h
include/sfs-0.6/mount_prot.x
include/sfs-0.6/msb.h
include/sfs-0.6/nfs3_ext.x
include/sfs-0.6/nfs3_nonnul.h
include/sfs-0.6/nfs3_prot.h
include/sfs-0.6/nfs3_prot.x
include/sfs-0.6/nfs3close_prot.h
include/sfs-0.6/nfs3close_prot.x
include/sfs-0.6/nfs3exp_prot.h
include/sfs-0.6/nfs3exp_prot.x
include/sfs-0.6/nfs_prot.h
include/sfs-0.6/nfs_prot.x
include/sfs-0.6/nfsmounter.h
include/sfs-0.6/nfsmounter.x
include/sfs-0.6/nfsserv.h
include/sfs-0.6/nfstrans.h
include/sfs-0.6/opnew.h
include/sfs-0.6/parseopt.h
include/sfs-0.6/password.h
include/sfs-0.6/pcre.h
include/sfs-0.6/pmap_prot.h
include/sfs-0.6/pmap_prot.x
include/sfs-0.6/prime.h
include/sfs-0.6/prng.h
include/sfs-0.6/qhash.h
include/sfs-0.6/rabin.h
include/sfs-0.6/refcnt.h
include/sfs-0.6/rex.h
include/sfs-0.6/rex_prot.h
include/sfs-0.6/rex_prot.x
include/sfs-0.6/rpctypes.h
include/sfs-0.6/rwfd.h
include/sfs-0.6/rxx.h
include/sfs-0.6/seqno.h
include/sfs-0.6/serial.h
include/sfs-0.6/sfs-internal.h
include/sfs-0.6/sfs_prot.h
include/sfs-0.6/sfs_prot.x
include/sfs-0.6/sfsagent.h
include/sfs-0.6/sfsagent.x
include/sfs-0.6/sfsauth_prot.h
include/sfs-0.6/sfsauth_prot.x
include/sfs-0.6/sfscd_prot.h
include/sfs-0.6/sfscd_prot.x
include/sfs-0.6/sfsclient.h
include/sfs-0.6/sfsconnect.h
include/sfs-0.6/sfskeymisc.h
include/sfs-0.6/sfsmisc.h
include/sfs-0.6/sfsro_prot.h
include/sfs-0.6/sfsro_prot.x
include/sfs-0.6/sfsserv.h
include/sfs-0.6/sfstty.h
include/sfs-0.6/sha1.h
include/sfs-0.6/srp.h
include/sfs-0.6/stllike.h
include/sfs-0.6/str.h
include/sfs-0.6/suio++.h
include/sfs-0.6/sysconf.h
include/sfs-0.6/tiger.h
include/sfs-0.6/union.h
include/sfs-0.6/uvfstrans.h
include/sfs-0.6/vatmpl.h
include/sfs-0.6/vec.h
include/sfs-0.6/wmstr.h
include/sfs-0.6/xdr_suio.h
include/sfs-0.6/xdrmisc.h
include/sfs-0.6/xhinfo.h
include/sfs.h
@dirrm include/sfs-0.6
@unexec install-info --delete %D/info/sfs.info %D/info/dir
info/sfs.info
info/sfs.info-1
info/sfs.info-2
@exec install-info %D/info/sfs.info %D/info/dir
lib/libsfs.a
lib/sfs
lib/sfs-0.6/aiod
lib/sfs-0.6/connect
lib/sfs-0.6/libarpc.a
lib/sfs-0.6/libarpc.la
lib/sfs-0.6/libasync.a
lib/sfs-0.6/libasync.la
lib/sfs-0.6/libsfscrypt.a
lib/sfs-0.6/libsfscrypt.la
lib/sfs-0.6/libsfsmisc.a
lib/sfs-0.6/libsfsmisc.la
lib/sfs-0.6/libsvc.a
lib/sfs-0.6/libsvc.la
lib/sfs-0.6/listen
lib/sfs-0.6/mallock.o
lib/sfs-0.6/moduled
lib/sfs-0.6/newaid
lib/sfs-0.6/nfsmounter
lib/sfs-0.6/pathinfo
lib/sfs-0.6/proxy
lib/sfs-0.6/ptyd
lib/sfs-0.6/rexd
lib/sfs-0.6/sfsauthd
lib/sfs-0.6/sfsrwcd
lib/sfs-0.6/sfsrwsd
lib/sfs-0.6/suidconnect
lib/sfs-0.6/ttyd
lib/sfs-0.6/xfer
@dirrm lib/sfs-0.6
sbin/funmount
sbin/sfscd
sbin/sfssd
share/doc/sfs/WELCOME
share/doc/sfs/README
share/doc/sfs/README.packageblurb
@dirrm share/doc/sfs
share/sfs/sfs_config
share/sfs/sfs_srp_parms
share/sfs/sfsauthd_config
share/sfs/sfscd_config
share/sfs/sfssd_config
@dirrm share/sfs
