# security/sfs/pkg-plist
bin/dirsearch
bin/newaid
bin/rex
bin/rpcc
bin/sfsagent
bin/sfsdeclog
bin/sfskey
bin/sfsproac
bin/ssu
@unexec if [ -f %D/etc/rc.d/sfscd.sh ]; then cmp -s %D/etc/rc.d/sfscd.sh.sample %D/etc/rc.d/sfscd.sh && rm -f %D/etc/rc.d/sfscd.sh || echo "${PKG_PREFIX}/etc/rc.d/sfscd.sh is different than distribution sample, leaving for manual cleanup." | fmt; fi
etc/rc.d/sfscd.sh.sample
@unexec if [ -f %D/etc/rc.d/sfssd.sh ]; then cmp -s %D/etc/rc.d/sfssd.sh.sample %D/etc/rc.d/sfssd.sh && rm -f %D/etc/rc.d/sfssd.sh || echo "${PKG_PREFIX}/etc/rc.d/sfssd.sh is different than distribution sample, leaving for manual cleanup." | fmt; fi
etc/rc.d/sfssd.sh.sample
etc/sfs/sfsrwsd_config.sample
include/sfs
include/sfs.h
include/sfs-%%SFS_VERSION%%/aclnt.h
include/sfs-%%SFS_VERSION%%/aes.h
include/sfs-%%SFS_VERSION%%/afsnode.h
include/sfs-%%SFS_VERSION%%/agentconn.h
include/sfs-%%SFS_VERSION%%/agentmisc.h
include/sfs-%%SFS_VERSION%%/aiod.h
include/sfs-%%SFS_VERSION%%/aiod_prot.h
include/sfs-%%SFS_VERSION%%/aios.h
include/sfs-%%SFS_VERSION%%/amisc.h
include/sfs-%%SFS_VERSION%%/arc4.h
include/sfs-%%SFS_VERSION%%/arena.h
include/sfs-%%SFS_VERSION%%/arpc.h
include/sfs-%%SFS_VERSION%%/array.h
include/sfs-%%SFS_VERSION%%/asrv.h
include/sfs-%%SFS_VERSION%%/async.h
include/sfs-%%SFS_VERSION%%/autoconf.h
include/sfs-%%SFS_VERSION%%/axprt.h
include/sfs-%%SFS_VERSION%%/axprt_crypt.h
include/sfs-%%SFS_VERSION%%/backoff.h
include/sfs-%%SFS_VERSION%%/bbuddy.h
include/sfs-%%SFS_VERSION%%/bench.h
include/sfs-%%SFS_VERSION%%/bigint.h
include/sfs-%%SFS_VERSION%%/bitvec.h
include/sfs-%%SFS_VERSION%%/blowfish.h
include/sfs-%%SFS_VERSION%%/callback.h
include/sfs-%%SFS_VERSION%%/cbuf.h
include/sfs-%%SFS_VERSION%%/crypt.h
include/sfs-%%SFS_VERSION%%/crypt_prot.h
include/sfs-%%SFS_VERSION%%/crypt_prot.x
include/sfs-%%SFS_VERSION%%/crypthash.h
include/sfs-%%SFS_VERSION%%/dns.h
include/sfs-%%SFS_VERSION%%/dnsparse.h
include/sfs-%%SFS_VERSION%%/err.h
include/sfs-%%SFS_VERSION%%/esign.h
include/sfs-%%SFS_VERSION%%/fdlim.h
include/sfs-%%SFS_VERSION%%/getfh3.h
include/sfs-%%SFS_VERSION%%/hashcash.h
include/sfs-%%SFS_VERSION%%/ihash.h
include/sfs-%%SFS_VERSION%%/init.h
include/sfs-%%SFS_VERSION%%/itree.h
include/sfs-%%SFS_VERSION%%/keyfunc.h
include/sfs-%%SFS_VERSION%%/list.h
include/sfs-%%SFS_VERSION%%/modalg.h
include/sfs-%%SFS_VERSION%%/mount_prot.h
include/sfs-%%SFS_VERSION%%/mount_prot.x
include/sfs-%%SFS_VERSION%%/msb.h
include/sfs-%%SFS_VERSION%%/nfs3_ext.x
include/sfs-%%SFS_VERSION%%/nfs3_nonnul.h
include/sfs-%%SFS_VERSION%%/nfs3_prot.h
include/sfs-%%SFS_VERSION%%/nfs3_prot.x
include/sfs-%%SFS_VERSION%%/nfs3close_prot.h
include/sfs-%%SFS_VERSION%%/nfs3close_prot.x
include/sfs-%%SFS_VERSION%%/nfs3exp_prot.h
include/sfs-%%SFS_VERSION%%/nfs3exp_prot.x
include/sfs-%%SFS_VERSION%%/nfs_prot.h
include/sfs-%%SFS_VERSION%%/nfs_prot.x
include/sfs-%%SFS_VERSION%%/nfsmounter.h
include/sfs-%%SFS_VERSION%%/nfsmounter.x
include/sfs-%%SFS_VERSION%%/nfsserv.h
include/sfs-%%SFS_VERSION%%/nfstrans.h
include/sfs-%%SFS_VERSION%%/opnew.h
include/sfs-%%SFS_VERSION%%/parseopt.h
include/sfs-%%SFS_VERSION%%/password.h
include/sfs-%%SFS_VERSION%%/pcre.h
include/sfs-%%SFS_VERSION%%/pmap_prot.h
include/sfs-%%SFS_VERSION%%/pmap_prot.x
include/sfs-%%SFS_VERSION%%/prime.h
include/sfs-%%SFS_VERSION%%/prng.h
include/sfs-%%SFS_VERSION%%/qhash.h
include/sfs-%%SFS_VERSION%%/rabin.h
include/sfs-%%SFS_VERSION%%/refcnt.h
include/sfs-%%SFS_VERSION%%/rex.h
include/sfs-%%SFS_VERSION%%/rex_prot.h
include/sfs-%%SFS_VERSION%%/rex_prot.x
include/sfs-%%SFS_VERSION%%/rpctypes.h
include/sfs-%%SFS_VERSION%%/rwfd.h
include/sfs-%%SFS_VERSION%%/rxx.h
include/sfs-%%SFS_VERSION%%/schnorr.h
include/sfs-%%SFS_VERSION%%/seqno.h
include/sfs-%%SFS_VERSION%%/serial.h
include/sfs-%%SFS_VERSION%%/sfs-internal.h
include/sfs-%%SFS_VERSION%%/sfs_prot.h
include/sfs-%%SFS_VERSION%%/sfs_prot.x
include/sfs-%%SFS_VERSION%%/sfsagent.h
include/sfs-%%SFS_VERSION%%/sfsagent.x
include/sfs-%%SFS_VERSION%%/sfsauth_prot.h
include/sfs-%%SFS_VERSION%%/sfsauth_prot.x
include/sfs-%%SFS_VERSION%%/sfscd_prot.h
include/sfs-%%SFS_VERSION%%/sfscd_prot.x
include/sfs-%%SFS_VERSION%%/sfsclient.h
include/sfs-%%SFS_VERSION%%/sfsconnect.h
include/sfs-%%SFS_VERSION%%/sfscrypt.h
include/sfs-%%SFS_VERSION%%/sfsextauth.h
include/sfs-%%SFS_VERSION%%/sfskeymgr.h
include/sfs-%%SFS_VERSION%%/sfskeymisc.h
include/sfs-%%SFS_VERSION%%/sfsmisc.h
include/sfs-%%SFS_VERSION%%/sfsro_prot.h
include/sfs-%%SFS_VERSION%%/sfsro_prot.x
include/sfs-%%SFS_VERSION%%/sfsschnorr.h
include/sfs-%%SFS_VERSION%%/sfsserv.h
include/sfs-%%SFS_VERSION%%/sfstty.h
include/sfs-%%SFS_VERSION%%/sha1.h
include/sfs-%%SFS_VERSION%%/srp.h
include/sfs-%%SFS_VERSION%%/stllike.h
include/sfs-%%SFS_VERSION%%/str.h
include/sfs-%%SFS_VERSION%%/suio++.h
include/sfs-%%SFS_VERSION%%/sysconf.h
include/sfs-%%SFS_VERSION%%/tiger.h
include/sfs-%%SFS_VERSION%%/union.h
include/sfs-%%SFS_VERSION%%/uvfstrans.h
include/sfs-%%SFS_VERSION%%/vatmpl.h
include/sfs-%%SFS_VERSION%%/vec.h
include/sfs-%%SFS_VERSION%%/wmstr.h
include/sfs-%%SFS_VERSION%%/xdr_suio.h
include/sfs-%%SFS_VERSION%%/xdrmisc.h
include/sfs-%%SFS_VERSION%%/xhinfo.h
@dirrm include/sfs-%%SFS_VERSION%%
lib/libsfs.a
lib/sfs
lib/sfs-%%SFS_VERSION%%/aiod
lib/sfs-%%SFS_VERSION%%/connect
lib/sfs-%%SFS_VERSION%%/libarpc.a
lib/sfs-%%SFS_VERSION%%/libasync.a
lib/sfs-%%SFS_VERSION%%/libsfscrypt.a
lib/sfs-%%SFS_VERSION%%/libsfsmisc.a
lib/sfs-%%SFS_VERSION%%/libsvc.a
lib/sfs-%%SFS_VERSION%%/listen
lib/sfs-%%SFS_VERSION%%/mallock.o
lib/sfs-%%SFS_VERSION%%/moduled
lib/sfs-%%SFS_VERSION%%/newaid
lib/sfs-%%SFS_VERSION%%/nfsmounter
lib/sfs-%%SFS_VERSION%%/pathinfo
lib/sfs-%%SFS_VERSION%%/proxy
lib/sfs-%%SFS_VERSION%%/ptyd
lib/sfs-%%SFS_VERSION%%/rexd
lib/sfs-%%SFS_VERSION%%/sfsauthd
lib/sfs-%%SFS_VERSION%%/sfsrwcd
lib/sfs-%%SFS_VERSION%%/sfsrwsd
lib/sfs-%%SFS_VERSION%%/upgradedb.pl
@group sfs
lib/sfs-%%SFS_VERSION%%/suidconnect
@group
lib/sfs-%%SFS_VERSION%%/ttyd
lib/sfs-%%SFS_VERSION%%/xfer
@dirrm lib/sfs-%%SFS_VERSION%%
sbin/funmount
sbin/sfscd
sbin/sfssd
sbin/vidb
%%DOCSDIR%%/WELCOME
%%DOCSDIR%%/README
%%DOCSDIR%%/README.packageblurb
%%DOCSDIR%%/README.0.7-upgrade
@dirrm %%DOCSDIR%%
share/sfs/agentrc
share/sfs/sfs_config
share/sfs/sfs_srp_parms
share/sfs/sfsauthd_config
share/sfs/sfscd_config
share/sfs/sfssd_config
@dirrm share/sfs
@dirrm etc/sfs
