# irc/irc2dc/files/openrc-irc2dc.in
#!/sbin/openrc-run

name="irc2dc"
command="%%PREFIX%%/bin/irc2dc"
pidfile="/var/run/irc2dc.pid"
required_files="%%PREFIX%%/etc/irc2dc.conf"

depend() {
	use net
	keyword -shutdown
}
