# devel/fam/pkg-plist
bin/fam
etc/fam.conf
etc/rc.d/fam.sh.sample
include/fam.h
lib/fam/editconf.perl
lib/libfam.a
lib/libfam.so
lib/libfam.so.0
@dirrm lib/fam
@unexec rm -f %D/etc/rc.d/fam.sh || true
