# www/sams2/pkg-plist
bin/sams2daemon
bin/sams2redir
bin/sams_send_email
bin/samsparser
%%APACHE%%@sample %%APACHEETCDIR%%/doc4sams2.conf.sample
%%APACHE%%@sample %%APACHEETCDIR%%/sams2.conf.sample
etc/rc.d/sams2
@sample etc/sams2.conf.sample
lib/sams2/libfsusage.a
lib/sams2/libfsusage.so
lib/sams2/libfsusage.so.1
lib/sams2/libfsusage.so.1.0.0
lib/sams2/libloadavg.a
lib/sams2/libloadavg.so
lib/sams2/libloadavg.so.1
lib/sams2/libloadavg.so.1.0.0
%%WWWDIR%%/authclass.php
%%WWWDIR%%/bin/diskspace
%%WWWDIR%%/bin/freemem
%%WWWDIR%%/bin/freeswap
%%WWWDIR%%/bin/fsusage
%%WWWDIR%%/bin/getntlmgroupid
%%WWWDIR%%/bin/getntlmgroups
%%WWWDIR%%/bin/getntlmusers
%%WWWDIR%%/bin/getwbinfousers
%%WWWDIR%%/bin/ntlm_auth
%%WWWDIR%%/bin/samsdf
%%WWWDIR%%/bin/syncfromad
%%WWWDIR%%/bin/syncfromldap
%%WWWDIR%%/bin/syncfromntlm
%%WWWDIR%%/bin/testntlmusergroup
%%WWWDIR%%/bin/testwbinfopasswd
%%WWWDIR%%/bin/updateuserlist
%%WWWDIR%%/bin/uptime
%%WWWDIR%%/blocked.php
%%WWWDIR%%/config.php
%%WWWDIR%%/config.php-e
%%WWWDIR%%/createdb.php
%%WWWDIR%%/dateclass.php
%%WWWDIR%%/dbclass.php
%%WWWDIR%%/dbclass.php-e
%%WWWDIR%%/dbstruct.sql
%%WWWDIR%%/groupclass.php
%%WWWDIR%%/icon/classic/adir.gif
%%WWWDIR%%/icon/classic/admin_48.jpg
%%WWWDIR%%/icon/classic/alllog.gif
%%WWWDIR%%/icon/classic/allow_48.jpg
%%WWWDIR%%/icon/classic/ark_48.jpg
%%WWWDIR%%/icon/classic/auth.gif
%%WWWDIR%%/icon/classic/backup_48.jpg
%%WWWDIR%%/icon/classic/blank.gif
%%WWWDIR%%/icon/classic/cache_32.jpg
%%WWWDIR%%/icon/classic/cache_48.jpg
%%WWWDIR%%/icon/classic/chname_32.jpg
%%WWWDIR%%/icon/classic/chname_48.jpg
%%WWWDIR%%/icon/classic/clock.gif
%%WWWDIR%%/icon/classic/clock_32.jpg
%%WWWDIR%%/icon/classic/clock_48.jpg
%%WWWDIR%%/icon/classic/config_20.jpg
%%WWWDIR%%/icon/classic/config_32.jpg
%%WWWDIR%%/icon/classic/config_48.jpg
%%WWWDIR%%/icon/classic/db.gif
%%WWWDIR%%/icon/classic/db_48.jpg
%%WWWDIR%%/icon/classic/dbcheck_32.jpg
%%WWWDIR%%/icon/classic/dbcheck_48.jpg
%%WWWDIR%%/icon/classic/dbupgrade_32.jpg
%%WWWDIR%%/icon/classic/dbupgrade_48.jpg
%%WWWDIR%%/icon/classic/delaypool.png
%%WWWDIR%%/icon/classic/delaypool_32.png
%%WWWDIR%%/icon/classic/denied.gif
%%WWWDIR%%/icon/classic/denied_48.jpg
%%WWWDIR%%/icon/classic/denied_big.gif
%%WWWDIR%%/icon/classic/diffDoc.gif
%%WWWDIR%%/icon/classic/diffFolder.gif
%%WWWDIR%%/icon/classic/domain-32.jpg
%%WWWDIR%%/icon/classic/domain-48.jpg
%%WWWDIR%%/icon/classic/domain_32.jpg
%%WWWDIR%%/icon/classic/dpersonal.gif
%%WWWDIR%%/icon/classic/duserd.gif
%%WWWDIR%%/icon/classic/earth.gif
%%WWWDIR%%/icon/classic/empty.gif
%%WWWDIR%%/icon/classic/erase_32.jpg
%%WWWDIR%%/icon/classic/erase_48.jpg
%%WWWDIR%%/icon/classic/error.gif
%%WWWDIR%%/icon/classic/export.jpg
%%WWWDIR%%/icon/classic/export_32.jpg
%%WWWDIR%%/icon/classic/export_48.jpg
%%WWWDIR%%/icon/classic/file.gif
%%WWWDIR%%/icon/classic/filesize_32.jpg
%%WWWDIR%%/icon/classic/filesize_48.jpg
%%WWWDIR%%/icon/classic/floppy.gif
%%WWWDIR%%/icon/classic/floppy_big.gif
%%WWWDIR%%/icon/classic/folder-closed.gif
%%WWWDIR%%/icon/classic/folder.gif
%%WWWDIR%%/icon/classic/ftraffic_32.jpg
%%WWWDIR%%/icon/classic/ftraffic_48.jpg
%%WWWDIR%%/icon/classic/ftv2blank.gif
%%WWWDIR%%/icon/classic/ftv2doc.gif
%%WWWDIR%%/icon/classic/ftv2folderclosed.gif
%%WWWDIR%%/icon/classic/ftv2folderopen.gif
%%WWWDIR%%/icon/classic/ftv2lastnode.gif
%%WWWDIR%%/icon/classic/ftv2link.gif
%%WWWDIR%%/icon/classic/ftv2mlastnode.gif
%%WWWDIR%%/icon/classic/ftv2mnode.gif
%%WWWDIR%%/icon/classic/ftv2node.gif
%%WWWDIR%%/icon/classic/ftv2plastnode.gif
%%WWWDIR%%/icon/classic/ftv2pnode.gif
%%WWWDIR%%/icon/classic/ftv2vertline.gif
%%WWWDIR%%/icon/classic/galka.gif
%%WWWDIR%%/icon/classic/getpassword.jpg
%%WWWDIR%%/icon/classic/grouplog-32.gif
%%WWWDIR%%/icon/classic/grouplog.gif
%%WWWDIR%%/icon/classic/help.jpg
%%WWWDIR%%/icon/classic/ident.gif
%%WWWDIR%%/icon/classic/import.jpg
%%WWWDIR%%/icon/classic/import_32.jpg
%%WWWDIR%%/icon/classic/import_48.jpg
%%WWWDIR%%/icon/classic/importcachedb_32.jpg
%%WWWDIR%%/icon/classic/importcachedb_48.jpg
%%WWWDIR%%/icon/classic/importdb_32.jpg
%%WWWDIR%%/icon/classic/importdb_48.jpg
%%WWWDIR%%/icon/classic/importfromad-32.jpg
%%WWWDIR%%/icon/classic/importfromad-48.jpg
%%WWWDIR%%/icon/classic/importfromldap-32.jpg
%%WWWDIR%%/icon/classic/importfromldap-48.jpg
%%WWWDIR%%/icon/classic/input.jpg
%%WWWDIR%%/icon/classic/join.gif
%%WWWDIR%%/icon/classic/joinbottom.gif
%%WWWDIR%%/icon/classic/lframe.jpg
%%WWWDIR%%/icon/classic/line.gif
%%WWWDIR%%/icon/classic/liststat_32.jpg
%%WWWDIR%%/icon/classic/liststat_48.jpg
%%WWWDIR%%/icon/classic/loadbackup_48.jpg
%%WWWDIR%%/icon/classic/loadbase-32.jpg
%%WWWDIR%%/icon/classic/loadbase.jpg
%%WWWDIR%%/icon/classic/loadbase_32.jpg
%%WWWDIR%%/icon/classic/loadbase_48.jpg
%%WWWDIR%%/icon/classic/loading.gif
%%WWWDIR%%/icon/classic/loadncsa_32.jpg
%%WWWDIR%%/icon/classic/loadncsa_48.jpg
%%WWWDIR%%/icon/classic/loadusers_32.jpg
%%WWWDIR%%/icon/classic/loadusers_48.jpg
%%WWWDIR%%/icon/classic/log.gif
%%WWWDIR%%/icon/classic/log_48.jpg
%%WWWDIR%%/icon/classic/logoff.gif
%%WWWDIR%%/icon/classic/logoff_20.gif
%%WWWDIR%%/icon/classic/logoff_20.jpg
%%WWWDIR%%/icon/classic/logoff_32.jpg
%%WWWDIR%%/icon/classic/logoff_48.jpg
%%WWWDIR%%/icon/classic/logs-32.gif
%%WWWDIR%%/icon/classic/logs.gif
%%WWWDIR%%/icon/classic/logs_squid.gif
%%WWWDIR%%/icon/classic/logs_squid_32.gif
%%WWWDIR%%/icon/classic/minus.gif
%%WWWDIR%%/icon/classic/minusbottom.gif
%%WWWDIR%%/icon/classic/monit.gif
%%WWWDIR%%/icon/classic/output.jpg
%%WWWDIR%%/icon/classic/paddressbook.gif
%%WWWDIR%%/icon/classic/passwd.gif
%%WWWDIR%%/icon/classic/password.gif
%%WWWDIR%%/icon/classic/pdf_32.jpg
%%WWWDIR%%/icon/classic/pdf_48.jpg
%%WWWDIR%%/icon/classic/pdir.gif
%%WWWDIR%%/icon/classic/pdir_002.gif
%%WWWDIR%%/icon/classic/persent_32.jpg
%%WWWDIR%%/icon/classic/persent_48.jpg
%%WWWDIR%%/icon/classic/personal.gif
%%WWWDIR%%/icon/classic/pfile.gif
%%WWWDIR%%/icon/classic/pgroup.gif
%%WWWDIR%%/icon/classic/pgroup_open.gif
%%WWWDIR%%/icon/classic/plus.gif
%%WWWDIR%%/icon/classic/plusbottom.gif
%%WWWDIR%%/icon/classic/plustop.gif
%%WWWDIR%%/icon/classic/pobject.gif
%%WWWDIR%%/icon/classic/pobject_open.gif
%%WWWDIR%%/icon/classic/porganization.gif
%%WWWDIR%%/icon/classic/printer.gif
%%WWWDIR%%/icon/classic/proxy.gif
%%WWWDIR%%/icon/classic/proxy_big.gif
%%WWWDIR%%/icon/classic/proxyes_32.jpg
%%WWWDIR%%/icon/classic/proxyes_48.jpg
%%WWWDIR%%/icon/classic/psite.gif
%%WWWDIR%%/icon/classic/puser.gif
%%WWWDIR%%/icon/classic/puser_open.gif
%%WWWDIR%%/icon/classic/puserd.gif
%%WWWDIR%%/icon/classic/quote_alarm.gif
%%WWWDIR%%/icon/classic/reark_48.jpg
%%WWWDIR%%/icon/classic/recalc_32.jpg
%%WWWDIR%%/icon/classic/recalc_48.jpg
%%WWWDIR%%/icon/classic/reconfig_48.jpg
%%WWWDIR%%/icon/classic/recsquid_32.jpg
%%WWWDIR%%/icon/classic/recsquid_48.jpg
%%WWWDIR%%/icon/classic/redir.jpg
%%WWWDIR%%/icon/classic/redirect.gif
%%WWWDIR%%/icon/classic/redirect_48.jpg
%%WWWDIR%%/icon/classic/redirect_big.gif
%%WWWDIR%%/icon/classic/roll.gif
%%WWWDIR%%/icon/classic/sams.gif
%%WWWDIR%%/icon/classic/savebase-32.jpg
%%WWWDIR%%/icon/classic/savebase.jpg
%%WWWDIR%%/icon/classic/savebase_32.jpg
%%WWWDIR%%/icon/classic/savebase_48.jpg
%%WWWDIR%%/icon/classic/setup.css
%%WWWDIR%%/icon/classic/shablon.jpg
%%WWWDIR%%/icon/classic/shutdown_32.jpg
%%WWWDIR%%/icon/classic/shutdown_48.jpg
%%WWWDIR%%/icon/classic/squidlog-32.gif
%%WWWDIR%%/icon/classic/squidlog.gif
%%WWWDIR%%/icon/classic/stat_32.jpg
%%WWWDIR%%/icon/classic/stat_48.jpg
%%WWWDIR%%/icon/classic/stop.gif
%%WWWDIR%%/icon/classic/stop.jpg
%%WWWDIR%%/icon/classic/stop2.gif
%%WWWDIR%%/icon/classic/stop_128.jpg
%%WWWDIR%%/icon/classic/stop_32.jpg
%%WWWDIR%%/icon/classic/straffic_32.jpg
%%WWWDIR%%/icon/classic/straffic_48.jpg
%%WWWDIR%%/icon/classic/subst_22.png
%%WWWDIR%%/icon/classic/switch_32.jpg
%%WWWDIR%%/icon/classic/switch_48.jpg
%%WWWDIR%%/icon/classic/switchuser_32.jpg
%%WWWDIR%%/icon/classic/switchuser_48.jpg
%%WWWDIR%%/icon/classic/sysplug_16.png
%%WWWDIR%%/icon/classic/sysplug_24.png
%%WWWDIR%%/icon/classic/sysplug_32.png
%%WWWDIR%%/icon/classic/sysplug_64.png
%%WWWDIR%%/icon/classic/tr_input.gif
%%WWWDIR%%/icon/classic/tr_input_32.jpg
%%WWWDIR%%/icon/classic/tr_input_48.jpg
%%WWWDIR%%/icon/classic/tr_output.gif
%%WWWDIR%%/icon/classic/tr_output_32.jpg
%%WWWDIR%%/icon/classic/tr_output_48.jpg
%%WWWDIR%%/icon/classic/tr_user_input_32.jpg
%%WWWDIR%%/icon/classic/tr_user_input_48.jpg
%%WWWDIR%%/icon/classic/tr_user_output_32.jpg
%%WWWDIR%%/icon/classic/tr_user_output_48.jpg
%%WWWDIR%%/icon/classic/tr_userport_input_32.jpg
%%WWWDIR%%/icon/classic/tr_userport_input_48.jpg
%%WWWDIR%%/icon/classic/tr_userport_output_32.jpg
%%WWWDIR%%/icon/classic/tr_userport_output_48.jpg
%%WWWDIR%%/icon/classic/traffic-16.jpg
%%WWWDIR%%/icon/classic/traffic-16.png
%%WWWDIR%%/icon/classic/traffic-32.jpg
%%WWWDIR%%/icon/classic/traffic.gif
%%WWWDIR%%/icon/classic/traffic.jpg
%%WWWDIR%%/icon/classic/traffic_32.jpg
%%WWWDIR%%/icon/classic/traffic_48.jpg
%%WWWDIR%%/icon/classic/trafmon-32.jpg
%%WWWDIR%%/icon/classic/trafmon-48.jpg
%%WWWDIR%%/icon/classic/trafmon-c-32.jpg
%%WWWDIR%%/icon/classic/trafmon-c-48.jpg
%%WWWDIR%%/icon/classic/trash-16.png
%%WWWDIR%%/icon/classic/trash_32.jpg
%%WWWDIR%%/icon/classic/trash_48.jpg
%%WWWDIR%%/icon/classic/tree.css
%%WWWDIR%%/icon/classic/treeview-black-line.gif
%%WWWDIR%%/icon/classic/treeview-black.gif
%%WWWDIR%%/icon/classic/treeview-default-line.gif
%%WWWDIR%%/icon/classic/treeview-default.gif
%%WWWDIR%%/icon/classic/treeview-famfamfam-line.gif
%%WWWDIR%%/icon/classic/treeview-famfamfam.gif
%%WWWDIR%%/icon/classic/treeview-gray-line.gif
%%WWWDIR%%/icon/classic/treeview-gray.gif
%%WWWDIR%%/icon/classic/treeview-red-line.gif
%%WWWDIR%%/icon/classic/treeview-red.gif
%%WWWDIR%%/icon/classic/treeview.css
%%WWWDIR%%/icon/classic/ttraffic_32.jpg
%%WWWDIR%%/icon/classic/ttraffic_48.jpg
%%WWWDIR%%/icon/classic/user.jpg
%%WWWDIR%%/icon/classic/user_32.jpg
%%WWWDIR%%/icon/classic/user_48.jpg
%%WWWDIR%%/icon/classic/user_active.png
%%WWWDIR%%/icon/classic/user_inactive.png
%%WWWDIR%%/icon/classic/user_moved.png
%%WWWDIR%%/icon/classic/user_off.png
%%WWWDIR%%/icon/classic/useradd.jpg
%%WWWDIR%%/icon/classic/useradd_32.jpg
%%WWWDIR%%/icon/classic/useradd_48.jpg
%%WWWDIR%%/icon/classic/usergroup_32.jpg
%%WWWDIR%%/icon/classic/usergroup_48.jpg
%%WWWDIR%%/icon/classic/userlog-32.gif
%%WWWDIR%%/icon/classic/userlog.gif
%%WWWDIR%%/icon/classic/usermon-32.jpg
%%WWWDIR%%/icon/classic/usermon-48.jpg
%%WWWDIR%%/icon/classic/usermon-c-32.jpg
%%WWWDIR%%/icon/classic/usermon-c-48.jpg
%%WWWDIR%%/icon/classic/userpasswd_32.jpg
%%WWWDIR%%/icon/classic/userpasswd_48.jpg
%%WWWDIR%%/icon/classic/userstat_32.jpg
%%WWWDIR%%/icon/classic/userstat_48.jpg
%%WWWDIR%%/icon/classic/usertraffic_32.jpg
%%WWWDIR%%/icon/classic/usertraffic_48.jpg
%%WWWDIR%%/icon/classic/warning.jpg
%%WWWDIR%%/icon/classic/web.gif
%%WWWDIR%%/icon/classic/webinterface.gif
%%WWWDIR%%/index.html
%%WWWDIR%%/index.php
%%WWWDIR%%/lang/koi8r-2-utf8
%%WWWDIR%%/lang/koi8r-2-win1251
%%WWWDIR%%/lang/koi8r-2-win1251.php
%%WWWDIR%%/lang/lang.EN
%%WWWDIR%%/lang/lang.KOI8-R
%%WWWDIR%%/lang/lang.UTF8
%%WWWDIR%%/lang/lang.WIN1251
%%WWWDIR%%/lang/utf8-2-koi8-r
%%WWWDIR%%/lang/win1251-2-koi8r
%%WWWDIR%%/lframe.php
%%WWWDIR%%/lframe_0ident.php
%%WWWDIR%%/lframe_1webconf.php
%%WWWDIR%%/lframe_2sysplug.php
%%WWWDIR%%/lframe_sams.php
%%WWWDIR%%/lframe_users.php
%%WWWDIR%%/lib/chart.php
%%WWWDIR%%/lib/charts.class.php
%%WWWDIR%%/lib/jquery-1.2.6.js
%%WWWDIR%%/lib/jquery.cookie.js
%%WWWDIR%%/lib/jquery.dataTables.js
%%WWWDIR%%/lib/jquery.treeview.js
%%WWWDIR%%/lib/piegraph.class.php
%%WWWDIR%%/lib/rgb.php
%%WWWDIR%%/lib/treeview.php
%%WWWDIR%%/main.php
%%WWWDIR%%/menu/ftiens4.js
%%WWWDIR%%/menu/menu.js
%%WWWDIR%%/menu/script.js
%%WWWDIR%%/menu/ua.js
%%WWWDIR%%/pluginmanager.php
%%WWWDIR%%/plugins/configtray/buttons/configbuttom_1_prop.php
%%WWWDIR%%/plugins/configtray/buttons/configbuttom_3_import.php
%%WWWDIR%%/plugins/configtray/configtraydescr.php
%%WWWDIR%%/plugins/configtray/configtrayimpl.php
%%WWWDIR%%/plugins/squid/buttons/squidbuttom_4_addcache.php
%%WWWDIR%%/plugins/squid/squiddescr.php
%%WWWDIR%%/plugins/squid/squidimpl.php
%%WWWDIR%%/plugins/webconfig/buttons/webconfigbuttom_1_prop.php
%%WWWDIR%%/plugins/webconfig/lang/lang.EN
%%WWWDIR%%/plugins/webconfig/webconfigdescr.php
%%WWWDIR%%/plugins/webconfig/webconfigimpl.php
%%WWWDIR%%/poolclass.php
%%WWWDIR%%/proxyclass.php
%%WWWDIR%%/reportsclass.php
%%WWWDIR%%/samsclass.php
%%WWWDIR%%/samstools.php
%%WWWDIR%%/setup.php
%%WWWDIR%%/shablonclass.php
%%WWWDIR%%/src/addurllist.php
%%WWWDIR%%/src/adldap.php
%%WWWDIR%%/src/adminbuttom_2_chuser.php
%%WWWDIR%%/src/adminbuttom_4_chpasswd.php
%%WWWDIR%%/src/admintray.php
%%WWWDIR%%/src/authadldbuttom_1_prop.php
%%WWWDIR%%/src/authadldbuttom_2_usersadd.php
%%WWWDIR%%/src/authadldbuttom_3_importfromad.php
%%WWWDIR%%/src/authadldtray.php
%%WWWDIR%%/src/authldapbuttom_1_prop.php
%%WWWDIR%%/src/authldapbuttom_2_usersadd.php
%%WWWDIR%%/src/authldapbuttom_3_importfromldap.php
%%WWWDIR%%/src/authldaptray.php
%%WWWDIR%%/src/authncsabuttom_1_usersimport.php
%%WWWDIR%%/src/authncsabuttom_2_ncsaimport.php
%%WWWDIR%%/src/authncsatray.php
%%WWWDIR%%/src/authntlmbuttom_1_prop.php
%%WWWDIR%%/src/authntlmbuttom_2_usersadd.php
%%WWWDIR%%/src/authntlmbuttom_3_importfromntlm.php
%%WWWDIR%%/src/authntlmtray.php
%%WWWDIR%%/src/authtray.php
%%WWWDIR%%/src/configbuttom_3_import.php
%%WWWDIR%%/src/configbuttom_3_importcache.php
%%WWWDIR%%/src/configbuttom_4_backup.php
%%WWWDIR%%/src/configbuttom_5_restore.php
%%WWWDIR%%/src/configbuttom_6_moveusertrafficdata.php
%%WWWDIR%%/src/configtray.php
%%WWWDIR%%/src/groupbuttom_0_prop.php
%%WWWDIR%%/src/groupbuttom_2_traffic.php
%%WWWDIR%%/src/groupbuttom_9_delete.php
%%WWWDIR%%/src/grouptray.php
%%WWWDIR%%/src/ldap.php
%%WWWDIR%%/src/lffolder_1_config.php
%%WWWDIR%%/src/lffolder_2_squid.php
%%WWWDIR%%/src/lffolder_3_auth.php
%%WWWDIR%%/src/lffolder_4_z0_filedenied.php
%%WWWDIR%%/src/lffolder_4_z0_redirect.php
%%WWWDIR%%/src/lffolder_4_z0_subst.php
%%WWWDIR%%/src/lffolder_4_z1_urldenied.php
%%WWWDIR%%/src/lffolder_4_z2_url.php
%%WWWDIR%%/src/lffolder_4_z3_urlallow.php
%%WWWDIR%%/src/lffolder_4_z4_local.php
%%WWWDIR%%/src/lffolder_5_z0_trange.php
%%WWWDIR%%/src/lffolder_5_z1_shablon.php
%%WWWDIR%%/src/lffolder_5_z20_delaypool.php
%%WWWDIR%%/src/lffolder_5_z2_groups.php
%%WWWDIR%%/src/newgrpbuttom_5_addgroup.php
%%WWWDIR%%/src/newgrptray.php
%%WWWDIR%%/src/poolbuttom_1_prop.php
%%WWWDIR%%/src/poolbuttom_9_delete.php
%%WWWDIR%%/src/pooltray.php
%%WWWDIR%%/src/proxybuttom_1_prop.php
%%WWWDIR%%/src/proxybuttom_2_reconfig.php
%%WWWDIR%%/src/proxybuttom_8_clear.php
%%WWWDIR%%/src/proxybuttom_9_delete.php
%%WWWDIR%%/src/proxybuttom_9_down.php
%%WWWDIR%%/src/proxytray.php
%%WWWDIR%%/src/redirbuttom_1_loadlist.php
%%WWWDIR%%/src/redirbuttom_2_savelist.php
%%WWWDIR%%/src/redirbuttom_3_deletelist.php
%%WWWDIR%%/src/redirlisttray.php
%%WWWDIR%%/src/shablonbuttom_1_prop.php
%%WWWDIR%%/src/shablonbuttom_9_delete.php
%%WWWDIR%%/src/shablonnew.php
%%WWWDIR%%/src/shablontray.php
%%WWWDIR%%/src/squidbuttom_4_addcache.php
%%WWWDIR%%/src/squidtray.php
%%WWWDIR%%/src/sysplugconfigtray.php
%%WWWDIR%%/src/trangebuttom_1_prop.php
%%WWWDIR%%/src/trangebuttom_9_delete.php
%%WWWDIR%%/src/trangetray.php
%%WWWDIR%%/src/userbuttom_1_prop.php
%%WWWDIR%%/src/userbuttom_2_traffic.php
%%WWWDIR%%/src/userbuttom_3_file.php
%%WWWDIR%%/src/userbuttom_3_time.php
%%WWWDIR%%/src/userbuttom_4_site.php
%%WWWDIR%%/src/userbuttom_5_clear.php
%%WWWDIR%%/src/userbuttom_6_chpasswd.php
%%WWWDIR%%/src/userbuttom_9_delete.php
%%WWWDIR%%/src/userbuttom_9_logoff.php
%%WWWDIR%%/src/usersbuttom_1_useradd.php
%%WWWDIR%%/src/usersbuttom_1_usersimport.php
%%WWWDIR%%/src/usersbuttom_2_traffic.php
%%WWWDIR%%/src/usersbuttom_3_site.php
%%WWWDIR%%/src/usersbuttom_4_file.php
%%WWWDIR%%/src/usersbuttom_4_percent.php
%%WWWDIR%%/src/usersbuttom_5_recalc.php
%%WWWDIR%%/src/usersbuttom_6_clear.php
%%WWWDIR%%/src/usersbuttom_91_test.php
%%WWWDIR%%/src/userstray.php
%%WWWDIR%%/src/usertray.php
%%WWWDIR%%/src/webconfigbuttom_1_prop.php
%%WWWDIR%%/src/webconfigtray.php
%%WWWDIR%%/test_table.html
%%WWWDIR%%/tools.php
%%WWWDIR%%/trangeclass.php
%%WWWDIR%%/tray.php
%%WWWDIR%%/userclass.php
