# irc/tr-ircd/pkg-plist
bin/ircd
bin/ircpasswd
bin/md5sum
etc/ircd.conf.sample
etc/ircd.links
etc/ircd.motd
etc/kline.conf
etc/rc.d/tr-ircd.sh.sample
include/blalloc.h
include/capabflag.h
include/chanmode.h
include/channel.h
include/class.h
include/client.h
include/cmodetab.h
include/common.h
include/comply.h
include/config.h
include/confitem.h
include/debug.h
include/defaults.h
include/dh.h
include/event.h
include/fd.h
include/fileio.h
include/find.h
include/floodspam.h
include/h.h
include/hash.h
include/hook.h
include/hostmask.h
include/identity.h
include/init.h
include/ircd_defs.h
include/ircsprintf.h
include/language.h
include/linebuf.h
include/listener.h
include/md5/closeout.h
include/md5/config.h
include/md5/error.h
include/md5/getline.h
include/md5/getopt.h
include/md5/gnukludge.h
include/md5/md5.h
include/md5/sys2.h
include/md5/system.h
include/md5/version-etc.h
include/md5/xalloc.h
include/memory.h
include/modulelist.h
include/modules.h
include/motditem.h
include/msg.h
include/msgtok1.h
include/numeric.h
include/packet.h
include/parse.h
include/patchlevel.h
include/protocol.h
include/protodef.h
include/queue.h
include/resnew.h
include/s_auth.h
include/s_bsd.h
include/s_conf.h
include/s_service.h
include/s_user.h
include/send.h
include/setup.h
include/struct.h
include/supported.h
include/sys.h
include/sysvar.h
include/throttle.h
include/tools.h
include/zlink.h
lib/chanmodes/cm_ban.so
lib/chanmodes/cm_chanban.so
lib/chanmodes/cm_chanop.so
lib/chanmodes/cm_excepts.so
lib/chanmodes/cm_flood.so
lib/chanmodes/cm_halfops.so
lib/chanmodes/cm_invites.so
lib/chanmodes/cm_key.so
lib/chanmodes/cm_limit.so
lib/chanmodes/cm_modhost.so
lib/chanmodes/cm_owner.so
lib/chanmodes/cm_protect.so
lib/chanmodes/cm_singles.so
lib/chanmodes/cm_voice.so
lib/contrib/m_admins.so
lib/contrib/m_getlang.so
lib/contrib/m_locops.so
lib/contrib/m_rnotice.so
lib/contrib/m_sadmins.so
lib/contrib/rehash_more.so
lib/contrib/wordfilter.so
lib/languages/lang_eng_funny.so
lib/languages/lang_eng_old.so
lib/languages/lang_turkish.so
lib/modules/m_admin.so
lib/modules/m_akill.so
lib/modules/m_away.so
lib/modules/m_burst.so
lib/modules/m_capab.so
lib/modules/m_close.so
lib/modules/m_connect.so
lib/modules/m_dccallow.so
lib/modules/m_die.so
lib/modules/m_display.so
lib/modules/m_dkey.so
lib/modules/m_error.so
lib/modules/m_exclude.so
lib/modules/m_globops.so
lib/modules/m_gnotice.so
lib/modules/m_goper.so
lib/modules/m_hash.so
lib/modules/m_help.so
lib/modules/m_info.so
lib/modules/m_invite.so
lib/modules/m_ircops.so
lib/modules/m_ison.so
lib/modules/m_join.so
lib/modules/m_jupiter.so
lib/modules/m_kick.so
lib/modules/m_kill.so
lib/modules/m_kline.so
lib/modules/m_links.so
lib/modules/m_list.so
lib/modules/m_lusers.so
lib/modules/m_map.so
lib/modules/m_message.so
lib/modules/m_mkpasswd.so
lib/modules/m_mode.so
lib/modules/m_motd.so
lib/modules/m_myid.so
lib/modules/m_names.so
lib/modules/m_nick.so
lib/modules/m_oper.so
lib/modules/m_operdo.so
lib/modules/m_part.so
lib/modules/m_pass.so
lib/modules/m_ping.so
lib/modules/m_pong.so
lib/modules/m_post.so
lib/modules/m_quit.so
lib/modules/m_rakill.so
lib/modules/m_rehash.so
lib/modules/m_restart.so
lib/modules/m_rexclude.so
lib/modules/m_rping.so
lib/modules/m_samode.so
lib/modules/m_server.so
lib/modules/m_service.so
lib/modules/m_services.so
lib/modules/m_servlist.so
lib/modules/m_servset.so
lib/modules/m_set.so
lib/modules/m_setlang.so
lib/modules/m_sgline.so
lib/modules/m_silence.so
lib/modules/m_sjoin.so
lib/modules/m_sqline.so
lib/modules/m_squery.so
lib/modules/m_squit.so
lib/modules/m_stats.so
lib/modules/m_svinfo.so
lib/modules/m_svsjoin.so
lib/modules/m_svskill.so
lib/modules/m_svsmode.so
lib/modules/m_svsnick.so
lib/modules/m_svsnoop.so
lib/modules/m_szline.so
lib/modules/m_time.so
lib/modules/m_topic.so
lib/modules/m_trace.so
lib/modules/m_unjupiter.so
lib/modules/m_unkline.so
lib/modules/m_unsgline.so
lib/modules/m_unsqline.so
lib/modules/m_unszline.so
lib/modules/m_user.so
lib/modules/m_userhost.so
lib/modules/m_userip.so
lib/modules/m_users.so
lib/modules/m_version.so
lib/modules/m_wallops.so
lib/modules/m_watch.so
lib/modules/m_who.so
lib/modules/m_whois.so
lib/modules/m_whowas.so
lib/modules/m_zline.so
lib/protocol/compat.so
lib/protocol/halcyon.so
lib/protocol/native.so
lib/protocol/pelennor.so
share/tr-ircd/help/en_us/admin
share/tr-ircd/help/en_us/akill
share/tr-ircd/help/en_us/away
share/tr-ircd/help/en_us/burst
share/tr-ircd/help/en_us/capab
share/tr-ircd/help/en_us/channelmodes
share/tr-ircd/help/en_us/close
share/tr-ircd/help/en_us/connect
share/tr-ircd/help/en_us/dccallow
share/tr-ircd/help/en_us/die
share/tr-ircd/help/en_us/display
share/tr-ircd/help/en_us/dkey
share/tr-ircd/help/en_us/error
share/tr-ircd/help/en_us/exclude
share/tr-ircd/help/en_us/gline
share/tr-ircd/help/en_us/globops
share/tr-ircd/help/en_us/gnotice
share/tr-ircd/help/en_us/goper
share/tr-ircd/help/en_us/hash
share/tr-ircd/help/en_us/help
share/tr-ircd/help/en_us/index
share/tr-ircd/help/en_us/info
share/tr-ircd/help/en_us/invite
share/tr-ircd/help/en_us/ircops
share/tr-ircd/help/en_us/ison
share/tr-ircd/help/en_us/join
share/tr-ircd/help/en_us/jupiter
share/tr-ircd/help/en_us/kick
share/tr-ircd/help/en_us/kill
share/tr-ircd/help/en_us/kline
share/tr-ircd/help/en_us/links
share/tr-ircd/help/en_us/list
share/tr-ircd/help/en_us/lusers
share/tr-ircd/help/en_us/map
share/tr-ircd/help/en_us/mkpasswd
share/tr-ircd/help/en_us/mode
share/tr-ircd/help/en_us/modules
share/tr-ircd/help/en_us/motd
share/tr-ircd/help/en_us/myid
share/tr-ircd/help/en_us/names
share/tr-ircd/help/en_us/netset
share/tr-ircd/help/en_us/nick
share/tr-ircd/help/en_us/notice
share/tr-ircd/help/en_us/oper
share/tr-ircd/help/en_us/operdo
share/tr-ircd/help/en_us/part
share/tr-ircd/help/en_us/pass
share/tr-ircd/help/en_us/ping
share/tr-ircd/help/en_us/pong
share/tr-ircd/help/en_us/privmsg
share/tr-ircd/help/en_us/quit
share/tr-ircd/help/en_us/rakill
share/tr-ircd/help/en_us/rehash
share/tr-ircd/help/en_us/restart
share/tr-ircd/help/en_us/rexclude
share/tr-ircd/help/en_us/samode
share/tr-ircd/help/en_us/server
share/tr-ircd/help/en_us/service
share/tr-ircd/help/en_us/services
share/tr-ircd/help/en_us/servlist
share/tr-ircd/help/en_us/servset
share/tr-ircd/help/en_us/set
share/tr-ircd/help/en_us/setlang
share/tr-ircd/help/en_us/sgline
share/tr-ircd/help/en_us/silence
share/tr-ircd/help/en_us/sjoin
share/tr-ircd/help/en_us/sqline
share/tr-ircd/help/en_us/squery
share/tr-ircd/help/en_us/squit
share/tr-ircd/help/en_us/stats
share/tr-ircd/help/en_us/svinfo
share/tr-ircd/help/en_us/svsjoin
share/tr-ircd/help/en_us/svskill
share/tr-ircd/help/en_us/svsmode
share/tr-ircd/help/en_us/svsnick
share/tr-ircd/help/en_us/svsnoop
share/tr-ircd/help/en_us/szline
share/tr-ircd/help/en_us/time
share/tr-ircd/help/en_us/topic
share/tr-ircd/help/en_us/trace
share/tr-ircd/help/en_us/ungline
share/tr-ircd/help/en_us/unjupiter
share/tr-ircd/help/en_us/unkline
share/tr-ircd/help/en_us/unsgline
share/tr-ircd/help/en_us/unsqline
share/tr-ircd/help/en_us/unszline
share/tr-ircd/help/en_us/user
share/tr-ircd/help/en_us/userhost
share/tr-ircd/help/en_us/userip
share/tr-ircd/help/en_us/usermodes
share/tr-ircd/help/en_us/users
share/tr-ircd/help/en_us/version
share/tr-ircd/help/en_us/wallops
share/tr-ircd/help/en_us/watch
share/tr-ircd/help/en_us/who
share/tr-ircd/help/en_us/whois
share/tr-ircd/help/en_us/whowas
share/tr-ircd/help/en_us/zline
@dirrm include/md5
@dirrm lib/chanmodes
@dirrm lib/contrib
@dirrm lib/languages
@dirrm lib/modules
@dirrm lib/protocol
@dirrm share/tr-ircd/help/en_us
@dirrm share/tr-ircd/help
@dirrm share/tr-ircd
