# irc/tr-ircd/pkg-plist
bin/converter.sh
bin/ircd
bin/ircpasswd
etc/ircd.motd
etc/kline.conf
etc/rc.d/tr-ircd.sh.sample
include/tr-ircd/md5/closeout.h
include/tr-ircd/md5/config.h
include/tr-ircd/md5/error.h
include/tr-ircd/md5/getline.h
include/tr-ircd/md5/getopt.h
include/tr-ircd/md5/gnukludge.h
include/tr-ircd/md5/md5.h
include/tr-ircd/md5/sys2.h
include/tr-ircd/md5/system.h
include/tr-ircd/md5/version-etc.h
include/tr-ircd/md5/xalloc.h
include/tr-ircd/blalloc.h
include/tr-ircd/capabflag.h
include/tr-ircd/chanmode.h
include/tr-ircd/channel.h
include/tr-ircd/class.h
include/tr-ircd/client.h
include/tr-ircd/cmodetab.h
include/tr-ircd/common.h
include/tr-ircd/comply.h
include/tr-ircd/config.h
include/tr-ircd/confitem.h
include/tr-ircd/debug.h
include/tr-ircd/defaults.h
include/tr-ircd/dh.h
include/tr-ircd/event.h
include/tr-ircd/fd.h
include/tr-ircd/fileio.h
include/tr-ircd/find.h
include/tr-ircd/floodspam.h
include/tr-ircd/h.h
include/tr-ircd/hash.h
include/tr-ircd/hook.h
include/tr-ircd/hostmask.h
include/tr-ircd/identity.h
include/tr-ircd/init.h
include/tr-ircd/ircd_defs.h
include/tr-ircd/ircsprintf.h
include/tr-ircd/language.h
include/tr-ircd/linebuf.h
include/tr-ircd/listener.h
include/tr-ircd/memory.h
include/tr-ircd/modulelist.h
include/tr-ircd/modules.h
include/tr-ircd/motditem.h
include/tr-ircd/msg.h
include/tr-ircd/msgtok1.h
include/tr-ircd/numeric.h
include/tr-ircd/packet.h
include/tr-ircd/parse.h
include/tr-ircd/patchlevel.h
include/tr-ircd/protocol.h
include/tr-ircd/protodef.h
include/tr-ircd/queue.h
include/tr-ircd/resnew.h
include/tr-ircd/s_auth.h
include/tr-ircd/s_bsd.h
include/tr-ircd/s_conf.h
include/tr-ircd/s_service.h
include/tr-ircd/s_user.h
include/tr-ircd/send.h
include/tr-ircd/struct.h
include/tr-ircd/supported.h
include/tr-ircd/sys.h
include/tr-ircd/sysvar.h
include/tr-ircd/throttle.h
include/tr-ircd/tools.h
include/tr-ircd/zlink.h
lib/tr-ircd/chanmodes/cm_ban.so
lib/tr-ircd/chanmodes/cm_chanban.so
lib/tr-ircd/chanmodes/cm_chanop.so
lib/tr-ircd/chanmodes/cm_excepts.so
lib/tr-ircd/chanmodes/cm_flood.so
lib/tr-ircd/chanmodes/cm_halfops.so
lib/tr-ircd/chanmodes/cm_invites.so
lib/tr-ircd/chanmodes/cm_key.so
lib/tr-ircd/chanmodes/cm_limit.so
lib/tr-ircd/chanmodes/cm_modhost.so
lib/tr-ircd/chanmodes/cm_owner.so
lib/tr-ircd/chanmodes/cm_protect.so
lib/tr-ircd/chanmodes/cm_singles.so
lib/tr-ircd/chanmodes/cm_voice.so
lib/tr-ircd/contrib/m_admins.so
lib/tr-ircd/contrib/m_locops.so
lib/tr-ircd/contrib/m_rnotice.so
lib/tr-ircd/contrib/m_sadmins.so
lib/tr-ircd/contrib/m_topic.so
lib/tr-ircd/contrib/rehash_more.so
lib/tr-ircd/contrib/wordfilter.so
lib/tr-ircd/languages/lang_eng_funny.so
lib/tr-ircd/languages/lang_eng_old.so
lib/tr-ircd/languages/lang_turkish.so
lib/tr-ircd/modules/m_admin.so
lib/tr-ircd/modules/m_akill.so
lib/tr-ircd/modules/m_away.so
lib/tr-ircd/modules/m_burst.so
lib/tr-ircd/modules/m_capab.so
lib/tr-ircd/modules/m_close.so
lib/tr-ircd/modules/m_connect.so
lib/tr-ircd/modules/m_dccallow.so
lib/tr-ircd/modules/m_die.so
lib/tr-ircd/modules/m_display.so
lib/tr-ircd/modules/m_dkey.so
lib/tr-ircd/modules/m_error.so
lib/tr-ircd/modules/m_exclude.so
lib/tr-ircd/modules/m_fakeip.so
lib/tr-ircd/modules/m_globops.so
lib/tr-ircd/modules/m_gnotice.so
lib/tr-ircd/modules/m_goper.so
lib/tr-ircd/modules/m_hash.so
lib/tr-ircd/modules/m_help.so
lib/tr-ircd/modules/m_info.so
lib/tr-ircd/modules/m_invite.so
lib/tr-ircd/modules/m_ircops.so
lib/tr-ircd/modules/m_ison.so
lib/tr-ircd/modules/m_join.so
lib/tr-ircd/modules/m_jupiter.so
lib/tr-ircd/modules/m_kick.so
lib/tr-ircd/modules/m_kill.so
lib/tr-ircd/modules/m_kline.so
lib/tr-ircd/modules/m_links.so
lib/tr-ircd/modules/m_list.so
lib/tr-ircd/modules/m_lusers.so
lib/tr-ircd/modules/m_map.so
lib/tr-ircd/modules/m_message.so
lib/tr-ircd/modules/m_mkpasswd.so
lib/tr-ircd/modules/m_mode.so
lib/tr-ircd/modules/m_motd.so
lib/tr-ircd/modules/m_myid.so
lib/tr-ircd/modules/m_names.so
lib/tr-ircd/modules/m_nick.so
lib/tr-ircd/modules/m_oper.so
lib/tr-ircd/modules/m_operdo.so
lib/tr-ircd/modules/m_part.so
lib/tr-ircd/modules/m_pass.so
lib/tr-ircd/modules/m_ping.so
lib/tr-ircd/modules/m_pong.so
lib/tr-ircd/modules/m_post.so
lib/tr-ircd/modules/m_quit.so
lib/tr-ircd/modules/m_rakill.so
lib/tr-ircd/modules/m_rehash.so
lib/tr-ircd/modules/m_restart.so
lib/tr-ircd/modules/m_rexclude.so
lib/tr-ircd/modules/m_rping.so
lib/tr-ircd/modules/m_samode.so
lib/tr-ircd/modules/m_server.so
lib/tr-ircd/modules/m_service.so
lib/tr-ircd/modules/m_services.so
lib/tr-ircd/modules/m_servlist.so
lib/tr-ircd/modules/m_servset.so
lib/tr-ircd/modules/m_set.so
lib/tr-ircd/modules/m_setlang.so
lib/tr-ircd/modules/m_sgline.so
lib/tr-ircd/modules/m_silence.so
lib/tr-ircd/modules/m_sjoin.so
lib/tr-ircd/modules/m_sqline.so
lib/tr-ircd/modules/m_squery.so
lib/tr-ircd/modules/m_squit.so
lib/tr-ircd/modules/m_stats.so
lib/tr-ircd/modules/m_svinfo.so
lib/tr-ircd/modules/m_svsjoin.so
lib/tr-ircd/modules/m_svskill.so
lib/tr-ircd/modules/m_svsmode.so
lib/tr-ircd/modules/m_svsnick.so
lib/tr-ircd/modules/m_svsnoop.so
lib/tr-ircd/modules/m_szline.so
lib/tr-ircd/modules/m_time.so
lib/tr-ircd/modules/m_topic.so
lib/tr-ircd/modules/m_trace.so
lib/tr-ircd/modules/m_unjupiter.so
lib/tr-ircd/modules/m_unkline.so
lib/tr-ircd/modules/m_unsgline.so
lib/tr-ircd/modules/m_unsqline.so
lib/tr-ircd/modules/m_unszline.so
lib/tr-ircd/modules/m_user.so
lib/tr-ircd/modules/m_userhost.so
lib/tr-ircd/modules/m_userip.so
lib/tr-ircd/modules/m_users.so
lib/tr-ircd/modules/m_version.so
lib/tr-ircd/modules/m_wallops.so
lib/tr-ircd/modules/m_watch.so
lib/tr-ircd/modules/m_who.so
lib/tr-ircd/modules/m_whois.so
lib/tr-ircd/modules/m_whowas.so
lib/tr-ircd/protocol/compat.so
lib/tr-ircd/protocol/halcyon.so
lib/tr-ircd/protocol/native.so
lib/tr-ircd/protocol/pelennor.so
libexec/md5sum
share/doc/tr-ircd/ChangeLog
share/doc/tr-ircd/Changes
share/doc/tr-ircd/INSTALL
share/doc/tr-ircd/LICENSE
share/doc/tr-ircd/Manual
share/doc/tr-ircd/Operators
share/doc/tr-ircd/README.TSora
share/doc/tr-ircd/README.small_nets
share/doc/tr-ircd/TR-Ircd-Team
share/doc/tr-ircd/example.conf
share/doc/tr-ircd/operguide.txt
share/doc/tr-ircd/opermyth.txt
share/doc/tr-ircd/opers.txt
share/doc/tr-ircd/serverhide.txt
share/tr-ircd/help/en_us/admin
share/tr-ircd/help/en_us/akill
share/tr-ircd/help/en_us/away
share/tr-ircd/help/en_us/burst
share/tr-ircd/help/en_us/capab
share/tr-ircd/help/en_us/channelmodes
share/tr-ircd/help/en_us/close
share/tr-ircd/help/en_us/connect
share/tr-ircd/help/en_us/dccallow
share/tr-ircd/help/en_us/die
share/tr-ircd/help/en_us/display
share/tr-ircd/help/en_us/dkey
share/tr-ircd/help/en_us/error
share/tr-ircd/help/en_us/exclude
share/tr-ircd/help/en_us/gline
share/tr-ircd/help/en_us/globops
share/tr-ircd/help/en_us/gnotice
share/tr-ircd/help/en_us/goper
share/tr-ircd/help/en_us/hash
share/tr-ircd/help/en_us/help
share/tr-ircd/help/en_us/index
share/tr-ircd/help/en_us/info
share/tr-ircd/help/en_us/invite
share/tr-ircd/help/en_us/ircops
share/tr-ircd/help/en_us/ison
share/tr-ircd/help/en_us/join
share/tr-ircd/help/en_us/jupiter
share/tr-ircd/help/en_us/kick
share/tr-ircd/help/en_us/kill
share/tr-ircd/help/en_us/kline
share/tr-ircd/help/en_us/links
share/tr-ircd/help/en_us/list
share/tr-ircd/help/en_us/lusers
share/tr-ircd/help/en_us/map
share/tr-ircd/help/en_us/mkpasswd
share/tr-ircd/help/en_us/mode
share/tr-ircd/help/en_us/modules
share/tr-ircd/help/en_us/motd
share/tr-ircd/help/en_us/myid
share/tr-ircd/help/en_us/names
share/tr-ircd/help/en_us/netset
share/tr-ircd/help/en_us/nick
share/tr-ircd/help/en_us/notice
share/tr-ircd/help/en_us/oper
share/tr-ircd/help/en_us/operdo
share/tr-ircd/help/en_us/part
share/tr-ircd/help/en_us/pass
share/tr-ircd/help/en_us/ping
share/tr-ircd/help/en_us/pong
share/tr-ircd/help/en_us/privmsg
share/tr-ircd/help/en_us/quit
share/tr-ircd/help/en_us/rakill
share/tr-ircd/help/en_us/rehash
share/tr-ircd/help/en_us/restart
share/tr-ircd/help/en_us/rexclude
share/tr-ircd/help/en_us/samode
share/tr-ircd/help/en_us/server
share/tr-ircd/help/en_us/service
share/tr-ircd/help/en_us/services
share/tr-ircd/help/en_us/servlist
share/tr-ircd/help/en_us/servset
share/tr-ircd/help/en_us/set
share/tr-ircd/help/en_us/setlang
share/tr-ircd/help/en_us/sgline
share/tr-ircd/help/en_us/silence
share/tr-ircd/help/en_us/sjoin
share/tr-ircd/help/en_us/sqline
share/tr-ircd/help/en_us/squery
share/tr-ircd/help/en_us/squit
share/tr-ircd/help/en_us/stats
share/tr-ircd/help/en_us/svinfo
share/tr-ircd/help/en_us/svsjoin
share/tr-ircd/help/en_us/svskill
share/tr-ircd/help/en_us/svsmode
share/tr-ircd/help/en_us/svsnick
share/tr-ircd/help/en_us/svsnoop
share/tr-ircd/help/en_us/szline
share/tr-ircd/help/en_us/time
share/tr-ircd/help/en_us/topic
share/tr-ircd/help/en_us/trace
share/tr-ircd/help/en_us/ungline
share/tr-ircd/help/en_us/unjupiter
share/tr-ircd/help/en_us/unkline
share/tr-ircd/help/en_us/unsgline
share/tr-ircd/help/en_us/unsqline
share/tr-ircd/help/en_us/unszline
share/tr-ircd/help/en_us/user
share/tr-ircd/help/en_us/userhost
share/tr-ircd/help/en_us/userip
share/tr-ircd/help/en_us/usermodes
share/tr-ircd/help/en_us/users
share/tr-ircd/help/en_us/version
share/tr-ircd/help/en_us/wallops
share/tr-ircd/help/en_us/watch
share/tr-ircd/help/en_us/who
share/tr-ircd/help/en_us/whois
share/tr-ircd/help/en_us/whowas
@dirrm include/tr-ircd/md5
@dirrm include/tr-ircd
@dirrm lib/tr-ircd/chanmodes
@dirrm lib/tr-ircd/contrib
@dirrm lib/tr-ircd/languages
@dirrm lib/tr-ircd/modules
@dirrm lib/tr-ircd/protocol
@dirrm lib/tr-ircd
@dirrm share/doc/tr-ircd
@dirrm share/tr-ircd/help/en_us
@dirrm share/tr-ircd/help
@dirrm share/tr-ircd
