# devel/py-libiocage/pkg-plist
bin/ioc
etc/rc.d/ioc
%%PYTHON_SITELIBDIR%%/ioc-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/PKG-INFO
%%PYTHON_SITELIBDIR%%/ioc-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/SOURCES.txt
%%PYTHON_SITELIBDIR%%/ioc-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/dependency_links.txt
%%PYTHON_SITELIBDIR%%/ioc-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/entry_points.txt
%%PYTHON_SITELIBDIR%%/ioc-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/requires.txt
%%PYTHON_SITELIBDIR%%/ioc-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/top_level.txt
%%PYTHON_SITELIBDIR%%/ioc/__init__.py
%%PYTHON_SITELIBDIR%%/ioc/activate.py
%%PYTHON_SITELIBDIR%%/ioc/clone.py
%%PYTHON_SITELIBDIR%%/ioc/console.py
%%PYTHON_SITELIBDIR%%/ioc/create.py
%%PYTHON_SITELIBDIR%%/ioc/deactivate.py
%%PYTHON_SITELIBDIR%%/ioc/destroy.py
%%PYTHON_SITELIBDIR%%/ioc/exec.py
%%PYTHON_SITELIBDIR%%/ioc/export.py
%%PYTHON_SITELIBDIR%%/ioc/fetch.py
%%PYTHON_SITELIBDIR%%/ioc/fstab.py
%%PYTHON_SITELIBDIR%%/ioc/get.py
%%PYTHON_SITELIBDIR%%/ioc/import.py
%%PYTHON_SITELIBDIR%%/ioc/list.py
%%PYTHON_SITELIBDIR%%/ioc/migrate.py
%%PYTHON_SITELIBDIR%%/ioc/pkg.py
%%PYTHON_SITELIBDIR%%/ioc/promote.py
%%PYTHON_SITELIBDIR%%/ioc/provision.py
%%PYTHON_SITELIBDIR%%/ioc/rename.py
%%PYTHON_SITELIBDIR%%/ioc/restart.py
%%PYTHON_SITELIBDIR%%/ioc/set.py
%%PYTHON_SITELIBDIR%%/ioc/shared/__init__.py
%%PYTHON_SITELIBDIR%%/ioc/shared/click.py
%%PYTHON_SITELIBDIR%%/ioc/shared/jail.py
%%PYTHON_SITELIBDIR%%/ioc/shared/output.py
%%PYTHON_SITELIBDIR%%/ioc/snapshot.py
%%PYTHON_SITELIBDIR%%/ioc/start.py
%%PYTHON_SITELIBDIR%%/ioc/stop.py
%%PYTHON_SITELIBDIR%%/ioc/update.py
%%PYTHON_SITELIBDIR%%/iocage-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/PKG-INFO
%%PYTHON_SITELIBDIR%%/iocage-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/SOURCES.txt
%%PYTHON_SITELIBDIR%%/iocage-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/dependency_links.txt
%%PYTHON_SITELIBDIR%%/iocage-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/requires.txt
%%PYTHON_SITELIBDIR%%/iocage-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/top_level.txt
%%PYTHON_SITELIBDIR%%/iocage/BridgeInterface.py
%%PYTHON_SITELIBDIR%%/iocage/CommandQueue.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Dataset.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/BaseConfig.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/Defaults.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/File/Fstab.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/File/RCConf.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/File/SysctlConf.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/File/__init__.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/JailConfig.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/Properties/Addresses.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/Properties/Defaultrouter.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/Properties/Interfaces.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/Properties/Resolver.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/Properties/ResourceLimit.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/Properties/__init__.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Jail/__init__.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Prototype.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Type/JSON.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Type/UCL.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Type/ZFS.py
%%PYTHON_SITELIBDIR%%/iocage/Config/Type/__init__.py
%%PYTHON_SITELIBDIR%%/iocage/Config/__init__.py
%%PYTHON_SITELIBDIR%%/iocage/Datasets.py
%%PYTHON_SITELIBDIR%%/iocage/DevfsRules.py
%%PYTHON_SITELIBDIR%%/iocage/Distribution.py
%%PYTHON_SITELIBDIR%%/iocage/Filter.py
%%PYTHON_SITELIBDIR%%/iocage/Firewall.py
%%PYTHON_SITELIBDIR%%/iocage/Host.py
%%PYTHON_SITELIBDIR%%/iocage/IPAddress.py
%%PYTHON_SITELIBDIR%%/iocage/Jail.py
%%PYTHON_SITELIBDIR%%/iocage/JailState.py
%%PYTHON_SITELIBDIR%%/iocage/Jails.py
%%PYTHON_SITELIBDIR%%/iocage/LaunchableResource.py
%%PYTHON_SITELIBDIR%%/iocage/ListableResource.py
%%PYTHON_SITELIBDIR%%/iocage/Logger.py
%%PYTHON_SITELIBDIR%%/iocage/MacAddress.py
%%PYTHON_SITELIBDIR%%/iocage/Network.py
%%PYTHON_SITELIBDIR%%/iocage/NetworkInterface.py
%%PYTHON_SITELIBDIR%%/iocage/NullFSBasejailStorage.py
%%PYTHON_SITELIBDIR%%/iocage/Pkg.py
%%PYTHON_SITELIBDIR%%/iocage/Prompts.py
%%PYTHON_SITELIBDIR%%/iocage/Provisioning/__init__.py
%%PYTHON_SITELIBDIR%%/iocage/Provisioning/ix.py
%%PYTHON_SITELIBDIR%%/iocage/Release.py
%%PYTHON_SITELIBDIR%%/iocage/Releases.py
%%PYTHON_SITELIBDIR%%/iocage/Resource.py
%%PYTHON_SITELIBDIR%%/iocage/ResourceBackup.py
%%PYTHON_SITELIBDIR%%/iocage/ResourceSelector.py
%%PYTHON_SITELIBDIR%%/iocage/ResourceUpdater.py
%%PYTHON_SITELIBDIR%%/iocage/SecureTarfile.py
%%PYTHON_SITELIBDIR%%/iocage/StandaloneJailStorage.py
%%PYTHON_SITELIBDIR%%/iocage/Storage.py
%%PYTHON_SITELIBDIR%%/iocage/Types.py
%%PYTHON_SITELIBDIR%%/iocage/VersionedResource.py
%%PYTHON_SITELIBDIR%%/iocage/ZFS.py
%%PYTHON_SITELIBDIR%%/iocage/ZFSBasejailStorage.py
%%PYTHON_SITELIBDIR%%/iocage/ZFSShareStorage.py
%%PYTHON_SITELIBDIR%%/iocage/__init__.py
%%PYTHON_SITELIBDIR%%/iocage/errors.py
%%PYTHON_SITELIBDIR%%/iocage/events.py
%%PYTHON_SITELIBDIR%%/iocage/helpers.py
%%PYTHON_SITELIBDIR%%/iocage/helpers_object.py
