# net/nrpe/pkg-plist
sbin/nrpe
etc/nrpe.cfg-sample
etc/rc.d/nrpe.sh
