# security/amavisd/pkg-plist
%%PORTDOCS%%%%DOCSDIR%%/FAQ
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README.exim
%%PORTDOCS%%%%DOCSDIR%%/README.milter
%%PORTDOCS%%%%DOCSDIR%%/README.postfix
%%PORTDOCS%%%%DOCSDIR%%/README.qmail
%%PORTDOCS%%%%DOCSDIR%%/README.scanners
%%PORTDOCS%%%%DOCSDIR%%/README.sendmail
%%PORTDOCS%%%%DOCSDIR%%/amavis.html
%%PORTDOCS%%%%DOCSDIR%%/amavis.m4
%%PORTDOCS%%%%DOCSDIR%%/amavis.png
%%PORTDOCS%%%%DOCSDIR%%/amavis.txt
%%PORTDOCS%%@dirrm %%DOCSDIR%%
etc/amavisd.conf
sbin/%%SBIN_AMAVIS%%
sbin/amavisd
%%MILTER%%etc/rc.d/amavisd-milter.sh
etc/rc.d/amavisd.sh
