# mail/mailgraph/files/mailgraph.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: mailgraph
# REQUIRE: DAEMON
#
# Add the following line to /etc/rc.conf to enable mailgraph:
#
# mailgraph_enable="YES"
#

. /etc/rc.subr

name="mailgraph"
rcvar=mailgraph_enable

command=%%PREFIX%%/sbin/${name}.pl
command_args='> /dev/null 2>&1'
command_interpreter=/usr/bin/perl
stop_postcmd=stop_postcmd

stop_postcmd()
{
    rm -f $pidfile
}

load_rc_config $name

: ${mailgraph_enable="NO"}
: ${mailgraph_pidfile="%%DATADIR%%/mailgraph.pid"}
: ${mailgraph_flags="--logfile /var/log/maillog --daemon-rrd=%%DATADIR%% --ignore-localhost --daemon --daemon-pid=${mailgraph_pidfile}"}
: ${mailgraph_user="%%MAILGRAPH_USER%%"}
: ${mailgraph_chdir="%%DATADIR%%"}

pidfile=${mailgraph_pidfile}

run_rc_command "$1"
