# mail/mailgraph/files/openrc-mailgraph.in
#!/sbin/openrc-run

name="mailgraph"
command=%%PREFIX%%/sbin/${name}.pl
command_args='> /dev/null 2>&1'

start_pre()
{
    mkdir -p "$mailgraph_chdir"
    chown "$mailgraph_user" "$mailgraph_chdir"
    if ! su -m "$mailgraph_user" -c "sh -c \"[ -r "$mailgraph_maillog" ]\""; then
        eerror "MAILGRAPH CANNOT RUN: $mailgraph_maillog not readable by $mailgraph_user"
        return 1
    fi
}

: ${mailgraph_chdir="%%DATADIR%%"}
: ${mailgraph_pidfile="$mailgraph_chdir/mailgraph.pid"}
: ${mailgraph_maillog="/var/log/maillog"}
: ${mailgraph_flags="--logfile $mailgraph_maillog --daemon-rrd=$mailgraph_chdir --ignore-localhost --daemon --daemon-pid=$mailgraph_pidfile"}
: ${mailgraph_user="%%MAILGRAPH_USER%%"}

pidfile=${mailgraph_pidfile}
