# sysutils/rsyslog8/files/rsyslogd.in
#!/bin/sh
#
# $FreeBSD$
#


# PROVIDE: rsyslogd
# REQUIRE: mountcritremote cleanvar newsyslog ldconfig
# BEFORE:  SERVERS

. /etc/rc.subr

name=rsyslogd
rcvar=rsyslogd_enable
command="%%PREFIX%%/sbin/${name}"
load_rc_config $name
: ${rsyslogd_enable:="NO"}
: ${rsyslogd_pidfile:="/var/run/rsyslogd.pid"}
: ${rsyslogd_config:="%%PREFIX%%/etc/rsyslog.conf"}
: ${rsyslogd_env:="TZ=/etc/localtime"}
pidfile="${rsyslogd_pidfile}"
command_args="-i ${pidfile} -f ${rsyslogd_config}"
required_files="${rsyslogd_config}"
extra_commands="reload"

run_rc_command "$1"
