# net/wackford-squeers/pkg-plist
bin/wackford
bin/squeers
etc/rc.d/squeers
