# comms/mwavem/pkg-plist
bin/mwavem
etc/mwavem.conf
etc/rc.d/mwavem.sh
share/mwavem/agc8.dsp
share/mwavem/agc9.dsp
share/mwavem/agc9wt.dsp
share/mwavem/async.dsp
share/mwavem/async2x.dsp
share/mwavem/c96441.dsp
share/mwavem/callerid.dsp
share/mwavem/callprog.dsp
share/mwavem/cancel.dsp
share/mwavem/cl1_hdlc.dsp
share/mwavem/class1.dsp
share/mwavem/class2.dsp
share/mwavem/class8.dsp
share/mwavem/class80.dsp
share/mwavem/disc.dsp
share/mwavem/dtmfr8.dsp
share/mwavem/dtmfrm.dsp
share/mwavem/fax01.dsp
share/mwavem/faxpmp17.dsp
share/mwavem/gain96.dsp
share/mwavem/gainwt.dsp
share/mwavem/gpc2pc.dsp
share/mwavem/gsm.dsp
share/mwavem/imaadpcm.dsp
share/mwavem/init0437.dsp
share/mwavem/int0896.dsp
share/mwavem/int9608.dsp
share/mwavem/meix0437.dsp
share/mwavem/mhint.dsp
share/mwavem/mixer96.dsp
share/mwavem/mnp.dsp
share/mwavem/mnp5.dsp
share/mwavem/modem.dsp
share/mwavem/mwbi0917.dsp
share/mwavem/mwos3780.dsp
share/mwavem/pc2gpc.dsp
share/mwavem/ramm.dsp
share/mwavem/rammv34.dsp
share/mwavem/router.dsp
share/mwavem/sildet.dsp
share/mwavem/spk076.dsp
share/mwavem/spk076m.dsp
share/mwavem/spkeec.dsp
share/mwavem/spkeecm.dsp
share/mwavem/spkmmtp.dsp
share/mwavem/spkvod.dsp
share/mwavem/stp08m.dsp
share/mwavem/stpvm.dsp
share/mwavem/str08ser.dsp
share/mwavem/strvm.dsp
share/mwavem/t30.dsp
share/mwavem/tonegen.dsp
share/mwavem/uartdma.dsp
share/mwavem/v22.dsp
share/mwavem/v23.dsp
share/mwavem/v32bis.dsp
share/mwavem/v34.dsp
share/mwavem/v34pcmc.dsp
share/mwavem/v42.dsp
share/mwavem/v42bisn.dsp
share/mwavem/v8.dsp
share/mwavem/v80.dsp
share/mwavem/v80mac.dsp
share/mwavem/v90.dsp
share/mwavem/vv_hdlc.dsp
share/mwavem/vv_prot.dsp
share/mwavem/mwmmem.prf
share/mwavem/tigr1040.wtt
share/mwavem/mwave.ko
@dirrm share/mwavem
@dirrm share/mwavem-fbsd
@exec rm -f /dev/mwave
@exec mknod /dev/mwave c 96 0
@unexec rm -f /dev/mwave
