# comms/mwavem/pkg-plist
sbin/mwavem
etc/mwavem.conf
etc/rc.d/mwavem.sh
share/mwave/agc8.dsp
share/mwave/agc9.dsp
share/mwave/agc9wt.dsp
share/mwave/async.dsp
share/mwave/async2x.dsp
share/mwave/c96441.dsp
share/mwave/callerid.dsp
share/mwave/callprog.dsp
share/mwave/cancel.dsp
share/mwave/cl1_hdlc.dsp
share/mwave/class1.dsp
share/mwave/class2.dsp
share/mwave/class8.dsp
share/mwave/class80.dsp
share/mwave/disc.dsp
share/mwave/dtmfr8.dsp
share/mwave/dtmfrm.dsp
share/mwave/fax01.dsp
share/mwave/faxpmp17.dsp
share/mwave/gain96.dsp
share/mwave/gainwt.dsp
share/mwave/gpc2pc.dsp
share/mwave/gsm.dsp
share/mwave/imaadpcm.dsp
share/mwave/init0437.dsp
share/mwave/int0896.dsp
share/mwave/int9608.dsp
share/mwave/meix0437.dsp
share/mwave/mhint.dsp
share/mwave/mixer96.dsp
share/mwave/mnp.dsp
share/mwave/mnp5.dsp
share/mwave/modem.dsp
share/mwave/mwbi0917.dsp
share/mwave/mwos3780.dsp
share/mwave/pc2gpc.dsp
share/mwave/ramm.dsp
share/mwave/rammv34.dsp
share/mwave/router.dsp
share/mwave/sildet.dsp
share/mwave/spk076.dsp
share/mwave/spk076m.dsp
share/mwave/spkeec.dsp
share/mwave/spkeecm.dsp
share/mwave/spkmmtp.dsp
share/mwave/spkvod.dsp
share/mwave/stp08m.dsp
share/mwave/stpvm.dsp
share/mwave/str08ser.dsp
share/mwave/strvm.dsp
share/mwave/t30.dsp
share/mwave/tonegen.dsp
share/mwave/uartdma.dsp
share/mwave/v22.dsp
share/mwave/v23.dsp
share/mwave/v32bis.dsp
share/mwave/v34.dsp
share/mwave/v34pcmc.dsp
share/mwave/v42.dsp
share/mwave/v42bisn.dsp
share/mwave/v8.dsp
share/mwave/v80.dsp
share/mwave/v80mac.dsp
share/mwave/v90.dsp
share/mwave/vv_hdlc.dsp
share/mwave/vv_prot.dsp
share/mwave/mwmmem.prf
share/mwave/tigr1040.wtt
share/mwave/mwavedd.ko
lib/libmwave.a
lib/libmeio.a
lib/libmwmutil.a
lib/libmwwtt32.a
lib/libmwmlw32.a
lib/libmwmpw32.a
lib/libmwmbl.a
@dirrm share/mwave
@exec rm -f /dev/mwave
@exec mknod /dev/mwave c 96 0
@unexec rm -f /dev/mwave
