# ftp/frox/files/frox.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: frox
# REQUIRE: NETWORKING
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable frox:
#
# frox_enable="YES"
# frox_flags="<set as needed>"
#
# See frox(8) for flags
#

. /etc/rc.subr

name=frox
rcvar=frox_enable

command=%%PREFIX%%/sbin/frox
pidfile=/var/run/frox.pid
required_files=%%PREFIX%%/etc/frox.conf

stop_postcmd=stop_postcmd

stop_postcmd()
{
   rm -f $pidfile
}

# set defaults

frox_enable=${frox_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"
