# sysutils/mgeupsd/pkg-plist
sbin/mgeupsd
sbin/powerhandler
etc/rc.d/mgeupsd.sh
man/man8/mgeupsd.8
