# net/openradius/files/openradius.in
#!/bin/sh

# PROVIDE: openradius
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable openradius:
#
# openradius_enable (bool):	Set to "YES" to enable openradius.
#				Default is "NO".
# openradius_flags (flags);	Set extra flags to openradius.
#				Default is "-o %%LOGFILE%%".
#

. /etc/rc.subr

name=openradius
rcvar=openradius_enable

load_rc_config $name

: ${openradius_enable:=NO}
: ${openradius_flags=-o %%LOGFILE%%}

command=%%PREFIX%%/sbin/radiusd

run_rc_command "$1"
