# sysutils/filewatcherd/files/openrc-filewatcherd.in
#!/sbin/openrc-run
#
# Add the following lines to /etc/rc.conf to enable filewatcherd:
#
# filewatcherd_conf (path):     Set full path to config file.
#                               Default is "%%PREFIX%%/etc/watchtab".

name=filewatcherd

: ${filewatcherd_conf="%%PREFIX%%/etc/watchtab"}

command=%%PREFIX%%/sbin/filewatcherd
command_args="${filewatcherd_conf}"
required_files=${filewatcherd_conf}

depend() {
	keyword -shutdown
}
