# net/tac_plus-libradius/files/tac_plus.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: tac_plus
# REQUIRE: DAEMON
# KEYWORD: shutdown

#
# Define these tac_plus_* variables in
#   /etc/rc.conf        or
#   /etc/rc.conf.local
#
# tac_plus_enabe="YES"
# tac_plus_cfg="%%PREFIX%%/etc/tac_plus/tac_plus.cfg"
# tac_plus_flags="-d 1" # see man (1) tac_plus
#

. /etc/rc.subr

name="tac_plus"
rcvar=`set_rcvar`

command=%%PREFIX%%/sbin/tac_plus
pidfile=/var/run/tac_plus.pid
tac_plus_enable=${tac_plus_enable:-"NO"}
tac_plus_flags=${tac_plus_flags:-"-d 1"}
tac_plus_configfile=${tac_plus_configfile:-"%%PREFIX%%/etc/tac_plus/tac_plus.cfg"}

load_rc_config $name

required_files=${tac_plus_configfile}
tac_plus_flags="-C ${tac_plus_configfile} ${tac_plus_flags}"

run_rc_command "$1"
