# mail/postpals/files/postpals.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: postpals
# REQUIRE: LOGIN cleanvar
# BEFORE: postpals-taillog 
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# postpals_enable="YES"
#

. /etc/rc.subr

name=postpals
rcvar=`set_rcvar`

command=%%PREFIX%%/sbin/${name}
pidfile="/var/run/${name}.pid"
command_interpreter="/usr/bin/perl"

stop_postcmd=stop_postcmd

stop_postcmd()
{
  rm -f $pidfile
}

load_rc_config $name

case "$postpals_enable" in
    [Yy][Ee][Ss] | 1 | [Oo][Nn] | [Tt][Rr][Uu][Ee]) ;;
    *) echo "To make use of $name you must first set $rcvar=\"YES\" in /etc/rc.conf" ;;
esac

: ${postpals_enable="NO"}
: ${postpals_tests="relay24,rcpt"}
: ${postpals_port="10040"}

command_args="-p ${pidfile} -u ${name} -g ${name} -d /var/db/${name}/${name}.db -T ${postpals_tests} ${postpals_port}"

run_rc_command "$1"
