# sysutils/moosefs-master/files/mfscgiserv.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: mfscgiserv
# REQUIRE: LOGIN 
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf or /etc/rc.conf.local to enable mfscgiserv:
#
# mfscgiserv_enable (bool):    Set to NO by default.
#                              Set it to YES to enable mfscgiserv.
# mfscgiserv_flags (str):      Set to "" by default.
#                              Extra flags passed to mfscgiserv.
#

. /etc/rc.subr

name=mfscgiserv
rcvar=mfscgiserv_enable

load_rc_config $name
: ${mfscgiserv_enable=NO}

mfscgiserv_user=mfs

command_interpreter=%%PYTHON_CMD%%
command=%%PREFIX%%/sbin/mfscgiserv

run_rc_command "$1"
