# sysutils/moosefs-master/files/mfsmetalogger.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: mfsmetalogger
# REQUIRE: LOGIN
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf to enable mfsmetalogger:
#
# mfsmetalogger_enable (bool):   Set to NO by default.
#                                Set it to YES to enable mfsmetalogger.
# mfsmetalogger_config (path):   Set to %%PREFIX%%/etc/mfs/mfsmetalogger.cfg
#                                by default.
# mfsmetalogger_flags (str):     Set to "" by default.
#                                Extra flags passed to mfsmetalogger.
#

. /etc/rc.subr

name=mfsmetalogger
rcvar=mfsmetalogger_enable

load_rc_config $name

: ${mfsmetalogger_enable=NO}
: ${mfsmetalogger_config=%%PREFIX%%/etc/mfs/mfsmetalogger.cfg}

extra_commands=reload
command=%%PREFIX%%/sbin/${name}
command_args="-c ${mfsmetalogger_config}"

run_rc_command "$1"
