# sysutils/tenshi/pkg-plist
bin/tenshi
etc/rc.d/tenshi.sh
etc/tenshi.conf.sample
