# mail/dma/files/openrc-dma_flushq.in
#!/sbin/openrc-run

name="dma_flushq"

start()
{
    %%PREFIX%%/libexec/dma -q
}

depend()
{
	keyword -shutdown
}
