# mail/quickml/pkg-plist
%%RUBY_LIBDIR%%/quickml/config.rb
%%RUBY_LIBDIR%%/quickml/gettext.rb
%%RUBY_LIBDIR%%/quickml/mail.rb
%%RUBY_LIBDIR%%/quickml/sweeper.rb
%%RUBY_LIBDIR%%/quickml/core.rb
%%RUBY_LIBDIR%%/quickml/logger.rb
%%RUBY_LIBDIR%%/quickml/server.rb
%%RUBY_LIBDIR%%/quickml/utils.rb
%%RUBY_LIBDIR%%/quickml/version.rb
%%RUBY_LIBDIR%%/quickml.rb
bin/quickml-analog
sbin/quickml
sbin/quickml-ctl
@unexec if cmp -s %D/etc/quickmlrc %D/etc/quickmlrc.sample; then rm -f %D/etc/quickmlrc; fi
etc/quickmlrc.sample
etc/rc.d/quickml.sh
%%DATADIR%%/messages.ja
@dirrm %%DATADIR%%
@dirrm %%RUBY_LIBDIR%%/quickml
%%PORTDOCS%%%%DOCSDIR%%/ml-usage.en.rd
%%PORTDOCS%%%%DOCSDIR%%/quickml.en.rd
%%PORTDOCS%%%%DOCSDIR%%/quickml.ja.rd
%%PORTDOCS%%%%DOCSDIR%%/with-mta.en.rd
%%PORTDOCS%%%%DOCSDIR%%/with-mta.ja.rd
%%PORTDOCS%%@dirrm %%DOCSDIR%%
