# sysutils/znapzend/files/znapzend.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: znapzend
# REQUIRE: LOGIN
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable znapzend:
# znapzend_enable (bool):	Set to "NO" by default.
#				Set it to "YES" to enable znapzend.
#

. /etc/rc.subr

name=znapzend
desc="Znapzend backup daemon"
rcvar=znapzend_enable

load_rc_config $name

: ${znapzend_enable:=NO}

extra_commands=reload
command_interpreter=/usr/local/bin/perl
sig_reload=HUP
pidfile=/var/run/${name}.pid
command=%%PREFIX%%/bin/${name}
command_args="--daemonize --pidfile=${pidfile} --logto=/var/log/znapzend.log"

run_rc_command "$1"
