# audio/mpdas/files/openrc-mpdas.in
#!/sbin/openrc-run

name=mpdas

command=%%PREFIX%%/bin/mpdas
command_args="-d"

depend() {
	need musicpd
	keyword -shutdown
}
