# net/amneziawg-tools/files/amneziawg.in
#!/bin/sh

# PROVIDE: amneziawg
# REQUIRE: NETWORKING
# KEYWORD: shutdown
#
# amneziawg_enable (bool):    Set to "YES" to enable amneziawg.
#                             (default: "NO")
#
# amneziawg_interfaces (str): List of interfaces to bring up/down
#                             on start/stop. (eg: "awg0 awg1")
#                             (default: "")
# amneziawg_env (str):        Environment variables for the userspace
#                             implementation. (eg: "LOG_LEVEL=debug")

. /etc/rc.subr

name=amneziawg
rcvar=amneziawg_enable
extra_commands="reload status"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
reload_cmd="${name}_reload"
status_cmd="${name}_status"

amneziawg_start()
{
	${amneziawg_env:+eval export $amneziawg_env}
    kldload -n if_awg

	for interface in ${amneziawg_interfaces}; do
		%%PREFIX%%/bin/awg-quick up ${interface}
	done
}

amneziawg_stop()
{
	for interface in ${amneziawg_interfaces}; do
		%%PREFIX%%/bin/awg-quick down ${interface}
	done
}

amneziawg_reload()
{
	${amneziawg_env:+eval export $amneziawg_env}

	for interface in ${amneziawg_interfaces}; do
		tmpfile="`mktemp`"
		%%PREFIX%%/bin/awg-quick strip ${interface} > ${tmpfile}
		%%PREFIX%%/bin/awg syncconf ${interface} ${tmpfile}
		rm -f ${tmpfile}
	done
}

amneziawg_status()
{
	${amneziawg_env:+eval export $amneziawg_env}
	amneziawg_status="0"

	for interface in ${amneziawg_interfaces}; do
		%%PREFIX%%/bin/awg show ${interface} || amneziawg_status="1"
	done

	return ${amneziawg_status}
}

load_rc_config $name

: ${amneziawg_enable="NO"}
: ${amneziawg_interfaces=""}
: ${amneziawg_env=""}

run_rc_command "$1"
