# net/crtmpserver/files/crtmpserver.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: crtmpserver
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# crtmpserver_enable (bool):   Set to NO by default.
#               Set it to YES to enable crtmpserver.
# crtmpserver_config (path):   Set to %%PREFIX%%/etc/crtmpserver.lua
#               by default.
# crtmpserver_username (username): Set to crtmpserver
#               by default.

. /etc/rc.subr

name="crtmpserver"
rcvar=${name}_enable

command=%%PREFIX%%/sbin/${name}

load_rc_config $name

: ${crtmpserver_enable="NO"}
: ${crtmpserver_config="%%PREFIX%%/etc/crtmpserver.lua"}
: ${crtmpserver_username="crtmpserver"}
: ${crtmpserver_groupname="crtmpserver"}
: ${crtmpserver_pidfile="/var/run/${name}.pid"}

required_files="$crtmpserver_config"

command_args="--daemon --uid=`id -u ${crtmpserver_username}` --gid=`pw groupshow ${crtmpserver_groupname} | awk -F ":" '{print $3}'` --pid=${crtmpserver_pidfile} $required_files"

run_rc_command "$1"

