# www/boa/pkg-plist
etc/rc.d/boa.sh.sample
sbin/boa
sbin/boa_indexer
%%EXAMPLESDIR%%/boa.conf
%%EXAMPLESDIR%%/cgi-test.cgi
%%EXAMPLESDIR%%/nph-test.cgi
%%EXAMPLESDIR%%/resolver.pl
@dirrm %%EXAMPLESDIR%%
