# ftp/ftp-proxy/files/ftp-proxy.in
#!/bin/sh
#
# PROVIDE: ftp-proxy
# REQUIRE: DAEMON pf
#
# Add the following lines to /etc/rc.conf to enable ftp-proxy:
#
# ftpproxy_enable (bool)	Set to "YES" to enable ftp-proxy.
#				Default is "NO".
# ftpproxy_flags (flags):	Set extra flags to ftp-proxy.
#				Default is "". See ftp-proxy(8).
#

. /etc/rc.subr

name="ftpproxy"
rcvar=${name}_enable

load_rc_config $name

: ${ftpproxy_enable="NO"}

command="%%PREFIX%%/sbin/ftp-proxy"

run_rc_command "$1"
