# databases/flare/files/flarei.sh.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: flarei
# REQUIRE: LOGIN
# BEFORE: flared
#
flarei_enable=${flarei_enable-"NO"}
flarei_config=${flarei_conffile-"%%PREFIX%%/etc/flarei.conf"}
flarei_pidfile=${flarei_pidfile-"/var/run/flarei.pid"}
flarei_flags="--daemonize"
flarei_sleepwait="2"

. /etc/rc.subr

name=flarei
rcvar=`set_rcvar`
command=%%PREFIX%%/bin/${name}

extra_commands="reload"
start_postcmd="flarei_poststart"

flarei_poststart () {
	sleep "${flarei_sleepwait}"
}

load_rc_config ${name}

case "${flarei_flags}" in
*-p\ *)
	echo "Warning: \$flarei_flags includes -p option." \
	    "Please use \$flarei_pidfile instead."
	;;
*--pid\ *)
	echo "Warning: \$flarei_flags includes -p option." \
	    "Please use \$flarei_pidfile instead."
	;;
*)
	flarei_flags="--pid ${flarei_pidfile} ${flarei_flags}"
	;;
esac

case "${flarei_flags}" in
*-f\ *)
	echo "Warning: \$flarei_flags includes -f option." \
	    "Please use \$flarei_config instead."
	;;
*--config\ *)
	echo "Warning: \$flarei_flags includes --config option." \
	    "Please use \$flarei_config instead."
	;;
*)
	flarei_flags="--config ${flarei_config} ${flarei_flags}"
	;;
esac

pidfile=${flarei_pidfile}
required_files=${flarei_conffile}

run_rc_command "$1"
