# net/apinger/files/openrc-apinger.in
#!/sbin/openrc-run

name=apinger
command=%%PREFIX%%/sbin/${name}
pidfile=/var/run/${name}.pid

depend() {
	keyword -shutdown
}
