# net/tac_plus-libradius/pkg-plist
sbin/tac_plus
etc/rc.d/tac_plus.sh.sample
etc/tac_plus/radius.example
etc/tac_plus/tac_plus.cfg
@dirrm etc/tac_plus
