# sysutils/relaunchd/pkg-plist
etc/rc.d/launchd
bin/launchctl
sbin/launchd
@dir /.launchd
@dir etc/launchd
@dir etc/launchd/agents
@dir etc/launchd/daemons
@dir share/launchd
@dir share/launchd/agents
@dir share/launchd/daemons
