# databases/riak/files/riak.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: riak
# REQUIRE: LOGIN cleanvar
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable riak:
# riak_enable (bool):		Set to "NO" by default.
#				Set it to "YES" to enable riak on boot.
#

. /etc/rc.subr

name="riak"
rcvar=riak_enable

pidfile=/var/run/riak/riak.pid

start_cmd="riak_start"
stop_cmd="riak_stop"
restart_cdm="riak_restart"
status_cmd="riak_status"
command="%%PREFIX%%/sbin/riak"

load_rc_config $name

# Read rc.d config and set defaults
load_rc_config "$name"
: ${riak_enable="NO"}

riak_start()
{
	echo "Starting Riak."
	%%PREFIX%%/sbin/riak start
	return 0
}

riak_stop()
{
	echo "Stopping Riak processes"
	%%PREFIX%%/sbin/riak stop
	killall -9 epmd
	return 0
}

riak_restart()
{
	riak_stop
	riak_start
	return 0
}

riak_status()
{
	if riak_running; then
		echo "Riak is running."
		return 0
	else
		echo "Riak is not running"
		return 1
	fi
}

riak_running()
{
	local pid result ps
	pid=`%%PREFIX%%/sbin/riak getpid`
	result=`echo $?`
	if [ "$result" == 0 ]; then
		ps=`ps -waux | grep ${pid} | grep riak`
		result=`echo $?`
		if [ "$result" ]; then
			return 0
		else
			return 1
		fi
	else
		return 1
	fi
}

run_rc_command "$1"
