# security/oidentd/pkg-plist
sbin/oidentd
etc/oidentd.conf.sample
etc/oidentd_masq.conf.sample
etc/rc.d/oidentd.sh
