# www/lycheeorg/files/lychee.cron.in
* * * * * www cd %%WWWDIR%% && php artisan schedule:run
