# net/rp-pppoe/pkg-plist
@postunexec for i in ppp/pppoe.conf rc.d/04pppoe-server.sh rc.d/04pppoe-relay.sh ; do if cmp -s %D/etc/$i.sample %D/etc/$i ; then rm %D/etc/$i ; fi ; done
etc/ppp/firewall-masq.sample
etc/ppp/firewall-standalone.sample
etc/ppp/pap-secrets.sample
etc/ppp/pppoe-server-options.sample
etc/ppp/pppoe.conf.sample
etc/rc.d/04pppoe-server.sh.sample
etc/rc.d/04pppoe-relay.sh.sample
man/man5/pppoe.conf.5.gz
man/man8/pppoe-connect.8.gz
man/man8/pppoe-relay.8.gz
man/man8/pppoe-server.8.gz
man/man8/pppoe-setup.8.gz
man/man8/pppoe-start.8.gz
man/man8/pppoe-status.8.gz
man/man8/pppoe-stop.8.gz
man/man8/pppoe.8.gz
sbin/pppoe-connect
sbin/pppoe-init
sbin/pppoe-setup
sbin/pppoe-start
sbin/pppoe-status
sbin/pppoe-stop
sbin/pppoe
sbin/pppoe-relay
sbin/pppoe-server
@postunexec rmdir %D/etc/ppp 2>/dev/null || echo "If you are permanently removing pppoe, you should also ``rm -Rf ${PKG_PREFIX}/etc/ppp'' removing any configuration and password files left." | /usr/bin/fmt
