# sysutils/usermin/pkg-plist
etc/rc.d/usermin.sh
@comment Insert PLIST.lib-usermin here
@comment -=[ begin PLIST.lib-usermin ]=-
@comment -=[ end PLIST.lib-usermin ]=-
@dirrm lib/usermin
