# comms/lirc/files/lircd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: lircd
# REQUIRE: DAEMON
#
# Add the following line to /etc/rc.conf to enable lircd:
#
# lircd_enable="YES"
#

. /etc/rc.subr

name="lircd"
rcvar=lircd_enable

load_rc_config ${name}
: ${lircd_enable="NO"}
: ${lircd_device="/dev/lirc0"}

procname=%%PREFIX%%/sbin/lircd
pidfile=/var/run/lirc/lircd.pid
lircd_config="%%PREFIX%%/etc/lircd.conf"
required_files=${lircd_config}
command="%%PREFIX%%/sbin/lircd"
command_args="-d ${lircd_device} ${lircd_config}"

run_rc_command "$1"
