# japanese/FreeWnn-server/files/twnn.in
#!/bin/sh
# $FreeBSD$
#
# PROVIDE: twnn
# REQUIRE: LOGIN
# KEYWORD: shutdown
#

. /etc/rc.subr

name=twnn
rcvar=twnn_enable
command=%%PREFIX%%/bin/cWnn4/tserver
procname=%%PREFIX%%/bin/cWnn4/tserver
start_precmd=start_precmd
stop_postcmd=stop_postcmd
WTCMD=%%PREFIX%%/bin/cWnn4/cwnntouch

start_precmd()
{
	rm -f /tmp/td_sockV4
	install -d -o wnn -m 755 %%PREFIX%%/lib/wnn/zh_TW/dic/usr
	${WTCMD} %%PREFIX%%/lib/wnn/zh_TW/dic/sys/*
}

stop_postcmd()
{
	rm -f /tmp/td_sockV4
}

load_rc_config ${name}
run_rc_command $1
