# net-im/jabber-pyicq/files/jabber-pyicq-transport.sh.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: jabber_pyicq
# REQUIRE: DAEMON %%JABBER_REQUIRE%%
# KEYWORD: shutdown
#
# Define these jabber_pyicq_* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#

. /etc/rc.subr

name="jabber_pyicq"
rcvar=jabber_pyicq_enable

load_rc_config $name

: ${jabber_pyicq_enable="NO"}
: ${jabber_pyicq_dir="%%PREFIX%%/lib/jabber/pyicq"}
: ${jabber_pyicq_piddir="%%JABBER_PIDDIR%%"}
: ${jabber_pyicq_user="%%JABBER_USER%%"}

pidfile="${jabber_pyicq_piddir}/PyICQt.pid"
command_interpreter="%%PYTHON_CMD%%"
command="${jabber_pyicq_dir}/PyICQt.py"
command_args="-b -o pid=${pidfile}"

run_rc_command "$1"
