# audio/musicpd/pkg-plist
bin/mpd
@unexec etc/rc.d/musicpd.sh stop 2>/dev/null || true
%%PORTDOCS%%%%DOCSDIR%%/COMMANDS
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/UPGRADING
%%PORTDOCS%%@dirrm %%DOCSDIR%%
%%EXAMPLESDIR%%/mpd.conf
@dirrm %%EXAMPLESDIR%%
etc/rc.d/musicpd.sh
