# net/freenet6/pkg-plist
bin/tspc
%%EXAMPLESDIR%%/tspc.conf.sample
bin/tspc-freebsd.sh
bin/checktunnel.sh
etc/rc.d/freenet6.sh
@dirrm %%EXAMPLESDIR%%
