# devel/zookeeper/files/zookeeper.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: zookeeper
#
# Add the following line to /etc/rc.conf to enable zookeeper:
#
# zookeeper_enable="YES"

command="/usr/local/zookeeper/bin/zkServer.sh"

. /etc/rc.subr

load_rc_config zookeeper

#set defaults
zookeeper_enable=${zookeeper_enable:-"NO"}

name=zookeeper
rcvar=zookeeper_enable

start_cmd="zookeeper_command start"
stop_cmd="zookeeper_command stop"
restart_cmd="zookeeper_command restart"
status_cmd="zookeeper_command status"

zookeeper_command()
{
	${command} ${rc_arg}
}

run_rc_command "$1"