# net/tac_plus4/pkg-plist
bin/tac_plus
bin/tac_pwd
etc/rc.d/tac_plus.sh
etc/tac_plus.conf.example
%%PORTDOCS%%%%DOCSDIR%%/convert.pl
%%PORTDOCS%%%%DOCSDIR%%/users_guide
%%PORTDOCS%%@dirrm %%DOCSDIR%%
