# net/tac_plus4/pkg-plist
bin/generate_passwd
etc/rc.d/tac_plus.sh
etc/tac_plus.conf.example
sbin/tac_plus
share/doc/tac_plus/users_guide
share/doc/tac_plus/convert.pl
@dirrm share/doc/tac_plus
