# www/hydra/files/hydra.in
#!/bin/sh
#
# PROVIDE: hydra
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# $FreeBSD$
#

hydra_enable=${hydra_enable-"NO"}

. /etc/rc.subr

load_rc_config hydra

name=hydra
rcvar=hydra_enable

command=%%PREFIX%%/bin/${name}

load_rc_config ${name}
run_rc_command "$1"
