# databases/neo4j/files/openrc-neo4j.in
#!/sbin/openrc-run

name=neo4j
description="Graphing database"

command=%%PREFIX%%/libexec/neo4j/neo4j

start_pre()
{
	NEO4J_DATA=%%DATADIR%%
	NEO4J_LIB=%%PREFIX%%/lib/neo4j
	NEO4J_LOGS=/var/log/neo4j
	NEO4J_PLUGINS=${NEO4J_LIB}/plugins
	NEO4J_RUN=/var/run/neo4j
	mkdir -p ${NEO4J_LOGS} ${NEO4J_RUN}
	$command start
}

stop_pre() {
	$command stop
}

status() {
	$command status
}

depend() {
	keyword -shutdown
}
