# news/diablo/pkg-plist
etc/rc.d/diablo.sh
news/dbin/dilookup
news/dbin/dicmd
news/dbin/didump
news/dbin/diconvhist
news/dbin/dexpire
news/dbin/didate
news/dbin/dspaminfo
news/dbin/diload
news/dbin/doutq
news/dbin/dclean
news/dbin/dspoolout
news/dbin/dreadart
news/dbin/dkp
news/dbin/pgpverify
news/dbin/dsyncgroups
news/dbin/dexpireover
news/dbin/diablo
news/dbin/dnewslink
news/dbin/plock
news/dbin/showlocks
news/dbin/dreaderd
news/samples/diablo.config
news/samples/rc.news
news/samples/dexpire.ctl
news/samples/dnewsfeeds
news/samples/dnntpspool.ctl
news/samples/diablo.hosts
news/samples/dactive.kp
news/samples/dreader.hosts
news/samples/dserver.hosts
news/samples/moderators
news/samples/dcontrol.ctl
news/dactive.kp
news/dcontrol.ctl
news/dexpire.ctl
news/diablo.config
news/diablo.hosts
news/dnewsfeeds
news/dnntpspool.ctl
news/dreader.hosts
news/dserver.hosts
news/moderators
news/rc.news
share/doc/diablo/COPYRIGHT
share/doc/diablo/INSTALL
share/doc/diablo/README.READER
share/doc/diablo/README.SERVER
share/doc/diablo/RELEASE_NOTES
share/doc/diablo/post-install-notes
@unexec /bin/rm -f /news
@dirrm share/doc/diablo
@dirrm news/dbin
@dirrm news/dqueue
@dirrm news/spool/news
@dirrm news/spool
@dirrm news/samples
@dirrm news
