# security/p5-openxpki-deployment/files/openxpki.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: openxpki
# REQUIRE: LOGIN
# REQUIRE: postgresql
# REQUIRE: mysql
# KEYWORD: shutdown

# Define this variable in the file /etc/rc.conf:
# openxpki_enable="YES"
# to ensure that openxpki starts at boot time.
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#
# openxpki_flags	(Set as needed)

openxpki_enable="${openxpki_enable-NO}"
openxpki_pidfile="%%PREFIX%%/var/openxpki/openxpki.pid"

case "$1" in 
    *start) command_args="start";;
esac

. /etc/rc.subr

name="openxpki"
rcvar=openxpki_enable
command="%%PREFIX%%/bin/openxpkictl"
procname="openxpkid"

load_rc_config $name

pidfile="${openxpki_pidfile}"

run_rc_command "$1"
