# net/icpld/pkg-plist
bin/icpld
etc/icpld.conf.sample
etc/rc.d/icpld.sh.sample
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%@dirrm %%DOCSDIR%%
