# devel/linux-f20-dbus/pkg-plist
bin/dbus-cleanup-sockets
bin/dbus-daemon
bin/dbus-monitor
bin/dbus-send
bin/dbus-uuidgen
etc/X11/xinit/xinitrc.d/00-start-message-bus.sh
%%ETCDIR%%-1/session.conf
%%ETCDIR%%-1/system.conf
lib/dbus-1/dbus-daemon-launch-helper
lib/systemd/system/dbus.service
lib/systemd/system/dbus.socket
lib/systemd/system/dbus.target.wants/dbus.socket
lib/systemd/system/messagebus.service
lib/systemd/system/multi-user.target.wants/dbus.service
lib/systemd/system/sockets.target.wants/dbus.socket
usr/bin/dbus-launch
%%PORTDOCS%%usr/%%DOCSDIR%%-1.6.12/COPYING
usr/share/man/man1/dbus-cleanup-sockets.1.gz
usr/share/man/man1/dbus-daemon.1.gz
usr/share/man/man1/dbus-launch.1.gz
usr/share/man/man1/dbus-monitor.1.gz
usr/share/man/man1/dbus-send.1.gz
usr/share/man/man1/dbus-uuidgen.1.gz
var/lib/dbus/machine-id
@dirrmtry bin
@dirrmtry etc/X11/xinit/xinitrc.d
@dirrmtry etc/X11/xinit
@dirrmtry %%ETCDIR%%-1/session.d
@dirrmtry %%ETCDIR%%-1/system.d
@dirrmtry %%ETCDIR%%-1
@dirrmtry lib/dbus-1
@dirrmtry lib/systemd/system/dbus.target.wants
@dirrmtry lib/systemd/system/multi-user.target.wants
@dirrmtry lib/systemd/system/sockets.target.wants
@dirrmtry lib/systemd/system
@dirrmtry lib/systemd
@dirrmtry lib
@dirrmtry usr/%%DATADIR%%-1/interfaces
@dirrmtry usr/%%DATADIR%%-1/services
@dirrmtry usr/%%DATADIR%%-1/system-services
@dirrmtry usr/%%DATADIR%%-1
%%PORTDOCS%%@dirrmtry usr/%%DOCSDIR%%-1.6.12
@dirrmtry var/lib/dbus
