# mail/py-spambayes/files/pyspamd.in
#!/bin/sh

# Start or stop pyspamd
# $FreeBSD$

# PROVIDE: pyspamd
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
prefix=/usr/local

# Define these pyspamd_* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#

. /etc/rc.subr

name="pyspamd"
rcvar=`set_rcvar`

load_rc_config $name

: ${pyspamd_enable="NO"}
: ${pyspamd_config="/usr/local/etc/bayescustomize.ini"}

command="/usr/sbin/daemon"
command_args=" -p /var/run/pyspamd.pid /usr/local/bin/python -OO /usr/local/bin/sb_server.py 2>&1 >/dev/null"
pidfile="/var/run/pyspamd.pid"
procname="/usr/local/bin/python"
export BAYESCUSTOMIZE=${pyspamd_config}

run_rc_command "$1"
