# net/kamailio/files/kamailio.in
#!/bin/sh
#
# $FreeBSD: net/kamailio/files/kamailio.in
#
# PROVIDE: kamailio
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable this service
# at system startup:
#
# kamailio_enable (bool): Set to NO by default.
# 			  Set it to YES to enable kamailio.
#

. /etc/rc.subr

name=kamailio
rcvar=kamailio_enable

load_rc_config ${name}

: ${kamailio_enable:=NO}
: ${kamailio_user:=%%USERS%%}
: ${pidfile:=%%RUNDIR%%/kamailio.pid}

command="/usr/local/sbin/kamailio"
command_args="-P ${pidfile}"

start_precmd="kamailio_prestart"

required_files="%%ETCDIR%%/kamailio.cfg"

kamailio_prestart() {
	install -d -m 0750 -o %%USERS%% -g %%GROUPS%% %%RUNDIR%%
}

run_rc_command $1
