# net/frp/files/frps.in
#!/bin/sh

# PROVIDE: frps
# REQUIRE: DAEMON NETWORKING
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable frps:
# frps_enable : set to "YES" to enable the daemon, default is "NO"

. /etc/rc.subr

name=frps
rcvar=frps_enable

load_rc_config $name

frps_enable=${frps_enable:-"NO"}

logfile="/var/log/${name}.log"

procname=%%PREFIX%%/bin/frps
command="/usr/sbin/daemon"
command_args="-u %%DAEMON_USER%% -o ${logfile} -t ${name} %%PREFIX%%/bin/frps -c %%PREFIX%%/etc/frps.ini"

run_rc_command "$1"
