# sysutils/lcdproc/files/openrc-lcdproc.in
#!/sbin/openrc-run

name="lcdproc"

: ${lcdproc_conf="%%PREFIX%%/etc/lcdproc.conf"}
: ${lcdproc_flags="-c ${lcdproc_conf}"}

required_files=${lcdproc_conf}
command="%%PREFIX%%/bin/lcdproc"
command_args=${lcdproc_flags}

depend() {
	need lcd
}
