# www/thundercache/files/thunder-rcNG.in
#!/bin/sh
#
# PROVIDE: thunder
# REQUIRE: LOGIN
#
# Remember to thunder_enable="YES" in /etc/rc.conf
#
# FreeBSD Brasil LTDA <consultoria@freebsdbrasil.com.br>
#

. /etc/rc.subr

name="thunder"
rcvar=`set_rcvar`
command="%%PREFIX%%/thunder/sbin/thunder ${thunder_flags}"
pidfile="/var/run/thunder/thunder.pid"
: ${thunder_enable="NO"}
: ${thunder_flags="-c %%PREFIX%%/thunder/etc/thunder.conf"}
load_rc_config $name
run_rc_command "$1"

