# www/xshttpd/files/openrc-xshttpd.in
#!/sbin/openrc-run

name="xshttpd"

pidfile="/var/run/httpd.pid"
command="%%PREFIX%%/bin/httpd"
