# net/dtcp/pkg-plist
sbin/dtcpc
sbin/dtcps
etc/rc.d/dtcpc.sh.sample
etc/rc.d/dtcps.sh.sample
