# net/dtcp/pkg-plist
sbin/dtcpauth
sbin/dtcpc
sbin/dtcps
etc/rc.d/dtcpc.sh
etc/rc.d/dtcps.sh
