# comms/pr/pkg-plist
etc/PortRedirector.conf.sample
etc/rc.d/pr.sh.sample
sbin/PortRedirector
sbin/PortTest
%%DOCSDIR%%/PR.html
@dirrm %%DOCSDIR%%
