# ftp/jftpgw/pkg-plist
sbin/jftpgw
etc/jftpgw.conf.sample
etc/rc.d/jftpgw.sh
