# net/openospfd/files/openospfd.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ospfd
# REQUIRE: LOGIN abi
# BEFORE:  securelevel
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable openospfd:
#
# openospfd_enable="YES"
# openospfd_flags="<set as needed>"
#
# See ospfd(8) for openospfd_flags

. /etc/rc.subr

name="openospfd"
rcvar=openospfd_enable

extra_commands="reload"
command="%%PREFIX%%/sbin/ospfd"
required_files="%%PREFIX%%/etc/ospfd.conf"

load_rc_config "$name"
: ${openospfd_enable="NO"}

run_rc_command "$1"
