# www/cherokee/files/cherokee.sh.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: cherokee
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable cherokee:
#
# cherokee_enable (bool):	Set it to "YES" to enable cherokee
#				Default is "NO".
# cherokee_conf (path):		Set full path to config file.
#				Default is "%%PREFIX%%/etc/cherokee/cherokee.conf".
#

. /etc/rc.subr

name="cherokee"
rcvar=cherokee_enable

load_rc_config $name

: ${cherokee_enable="NO"}
: ${cherokee_conf="%%PREFIX%%/etc/cherokee/cherokee.conf"}

command=%%PREFIX%%/sbin/cherokee
command_args="-C ${cherokee_conf} -d"
pidfile=/var/run/cherokee.pid
required_files=${cherokee_conf}
stop_postcmd=stop_postcmd

stop_postcmd()
{
	rm -f ${pidfile}
}

run_rc_command "$1"
