# sysutils/ganglia-monitor-core/files/gmond.sh.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: gmond
# REQUIRE: DAEMON
# KEYWORD: shutdown

bname=`basename $0`

. /etc/rc.subr

name=`echo "${bname}" | tr ".-" "__"`
rcvar=`set_rcvar`
command="%%PREFIX%%/sbin/gmond"

load_rc_config ganglia

load_rc_config $bname

eval "${name}_enable=\${${name}_enable-NO}"
eval "conffile=\${${name}_conf-%%PREFIX%%/etc/${bname}.conf}"

pidfile="/var/run/${bname}.pid"
command_args="-c ${conffile} -p ${pidfile}"
required_files=${conffile}

run_rc_command $*
