# sysutils/ganglia-monitor-core/files/openrc-gmond.in
#!/sbin/openrc-run

name=gmond

command="%%PREFIX%%/sbin/gmond"

eval "conffile=\${gmond_conf-%%PREFIX%%/etc/gmond.conf}"

pidfile="/var/run/${bname}.pid"
command_args="-c ${conffile} -p ${pidfile}"
required_files=${conffile}
