# net/mpd4/pkg-plist
etc/mpd4/mpd.conf.sample
etc/mpd4/mpd.links.sample
etc/mpd4/mpd.secret.sample
etc/mpd4/mpd.script.sample
etc/rc.d/mpd4.sh.sample
sbin/mpd4
share/doc/mpd4/README
share/doc/mpd4/mpd.ps
share/doc/mpd4/mpd.html
share/doc/mpd4/mpd1.html
share/doc/mpd4/mpd10.html
share/doc/mpd4/mpd11.html
share/doc/mpd4/mpd12.html
share/doc/mpd4/mpd13.html
share/doc/mpd4/mpd14.html
share/doc/mpd4/mpd15.html
share/doc/mpd4/mpd16.html
share/doc/mpd4/mpd17.html
share/doc/mpd4/mpd18.html
share/doc/mpd4/mpd19.html
share/doc/mpd4/mpd2.html
share/doc/mpd4/mpd20.html
share/doc/mpd4/mpd21.html
share/doc/mpd4/mpd22.html
share/doc/mpd4/mpd23.html
share/doc/mpd4/mpd24.html
share/doc/mpd4/mpd25.html
share/doc/mpd4/mpd26.html
share/doc/mpd4/mpd27.html
share/doc/mpd4/mpd28.html
share/doc/mpd4/mpd29.html
share/doc/mpd4/mpd3.html
share/doc/mpd4/mpd30.html
share/doc/mpd4/mpd31.html
share/doc/mpd4/mpd32.html
share/doc/mpd4/mpd33.html
share/doc/mpd4/mpd34.html
share/doc/mpd4/mpd35.html
share/doc/mpd4/mpd36.html
share/doc/mpd4/mpd37.html
share/doc/mpd4/mpd38.html
share/doc/mpd4/mpd39.html
share/doc/mpd4/mpd4.html
share/doc/mpd4/mpd40.html
share/doc/mpd4/mpd41.html
share/doc/mpd4/mpd42.html
share/doc/mpd4/mpd43.html
share/doc/mpd4/mpd44.html
share/doc/mpd4/mpd45.html
share/doc/mpd4/mpd46.html
share/doc/mpd4/mpd47.html
share/doc/mpd4/mpd48.html
share/doc/mpd4/mpd49.html
share/doc/mpd4/mpd5.html
share/doc/mpd4/mpd50.html
share/doc/mpd4/mpd51.html
share/doc/mpd4/mpd52.html
share/doc/mpd4/mpd53.html
share/doc/mpd4/mpd54.html
share/doc/mpd4/mpd55.html
share/doc/mpd4/mpd56.html
share/doc/mpd4/mpd57.html
share/doc/mpd4/mpd58.html
share/doc/mpd4/mpd59.html
share/doc/mpd4/mpd6.html
share/doc/mpd4/mpd7.html
share/doc/mpd4/mpd8.html
share/doc/mpd4/mpd9.html
share/doc/mpd4/mpd_toc.html
@dirrm etc/mpd4
@dirrm share/doc/mpd4
