# audio/liquidsoap/files/liquidsoap.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: liquidsoap
# REQUIRE: NETWORKING icecast2 shoutcast
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable liquidsoap:
#
# liquidsoap_enable="YES"
#

. /etc/rc.subr

name=liquidsoap
rcvar=liquidsoap_enable

command=%%PREFIX%%/bin/liquidsoap

load_rc_config $name

: ${liquidsoap_enable:="NO"}
: ${liquidsoap_pidfile:="%%PIDDIR%%/radio.pid"}
: ${liquidsoap_config:="%%ETCDIR%%/radio.liq"}
: ${liquidsoap_user:="%%LIQUIDSOAP_USER%%"}
: ${liquidsoap_group:="%%LIQUIDSOAP_GROUP%%"}

required_files="${liquidsoap_config}"
command_args="-q -d ${liquidsoap_config}"

run_rc_command "$1"
