# ftp/jftpgw/files/openrc-jftpgw.in
#!/sbin/openrc-run

name=jftpgw
command=%%PREFIX%%/sbin/jftpgw
jftpgw_config=${jftpgw_config:-"%%PREFIX%%/etc/jftpgw.conf"}
required_files=${jftpgw_config}
command_args="-f ${jftpgw_config}"

depend() {
	use net
	keyword -shutdown
}
