# net/xorp/pkg-plist
bin/call_xrl
bin/xorp_profiler
bin/xorp_rtrmgr
bin/xorpsh
bgp/xorp_bgp
bgp/tools/xorpsh_print_peers
bgp/tools/xorpsh_print_routes
cli/tools/send_cli_processor_xrl
contrib/olsr/tools/clear_database
contrib/olsr/tools/print_databases
contrib/olsr/xorp_olsr4
etc/rc.d/xorp
etc/templates/bgp.cmds
etc/templates/bgp.tp
etc/templates/fea.cmds
etc/templates/fea.tp
etc/templates/fib2mrib.tp
etc/templates/firewall.tp
etc/templates/host.cmds
etc/templates/igmp.cmds
etc/templates/igmp.tp
etc/templates/interfaces.tp
etc/templates/mfea.cmds
etc/templates/mfea4.tp
etc/templates/mfea6.cmds
etc/templates/mfea6.tp
etc/templates/misc.cmds
etc/templates/mld.cmds
etc/templates/mld.tp
etc/templates/olsr4.cmds
etc/templates/olsr4.tp
etc/templates/ospfv2.cmds
etc/templates/ospfv2.tp
etc/templates/ospfv3.cmds
etc/templates/ospfv3.tp
etc/templates/pim.cmds
etc/templates/pim6.cmds
etc/templates/rib.tp
etc/templates/pimsm4.tp
etc/templates/pimsm6.tp
etc/templates/plumbing.tp
etc/templates/policy.tp
etc/templates/protocols.tp
etc/templates/rib.cmds
etc/templates/rip.cmds
etc/templates/rip.tp
etc/templates/ripng.cmds
etc/templates/ripng.tp
etc/templates/rtrmgr.tp
etc/templates/snmp.tp
etc/templates/static_routes.tp
etc/templates/xorpsh.cmds
fea/xorp_fea
fea/xorp_fea_click_config_generator
fea/xorp_fea_dummy
fea/tools/show_interfaces
fib2mrib/xorp_fib2mrib
libxipc/xorp_finder
mld6igmp/xorp_igmp
mld6igmp/xorp_mld
ospf/xorp_ospfv2
ospf/xorp_ospfv3
ospf/tools/clear_database
ospf/tools/print_lsas
ospf/tools/print_neighbours
pim/xorp_pimsm4
pim/xorp_pimsm6
policy/xorp_policy
rib/xorp_rib
rib/tools/show_distances
rib/tools/show_routes
rip/xorp_rip
rip/xorp_ripng
rip/tools/show_peer_stats
rip/tools/show_stats
static_routes/xorp_static_routes
xrl/targets/bgp.xrls
xrl/targets/bgp4_mib.xrls
xrl/targets/cli.xrls
xrl/targets/coord.xrls
xrl/targets/fea.xrls
xrl/targets/fea_ifmgr_mirror.xrls
xrl/targets/fib2mrib.xrls
xrl/targets/finder.xrls
xrl/targets/finder_client.xrls
xrl/targets/mfea.xrls
xrl/targets/mld6igmp.xrls
xrl/targets/olsr4.xrls
xrl/targets/ospfv2.xrls
xrl/targets/ospfv3.xrls
xrl/targets/pim.xrls
xrl/targets/policy.xrls
xrl/targets/profiler.xrls
xrl/targets/rib.xrls
xrl/targets/ribclient.xrls
xrl/targets/rip.xrls
xrl/targets/ripng.xrls
xrl/targets/rtrmgr.xrls
xrl/targets/show_distances.xrls
xrl/targets/show_routes.xrls
xrl/targets/static_routes.xrls
xrl/targets/test.xrls
xrl/targets/test_fea_ifmgr_mirror.xrls
xrl/targets/test_fea_rawlink.xrls
xrl/targets/test_finder_events.xrls
xrl/targets/test_peer.xrls
xrl/targets/test_socket4.xrls
xrl/targets/test_socket6.xrls
xrl/targets/test_xrls.xrls
xrl/targets/xorp_if_mib.xrls
xrl/targets/xorpsh.xrls
%%XORP_SNMP%%libdata/ldconfig/xorp
%%XORP_SNMP%%mibs/bgp4_mib_1657.la
%%XORP_SNMP%%mibs/bgp4_mib_1657.so
%%XORP_SNMP%%mibs/libnetsnmpxorp.la
%%XORP_SNMP%%mibs/libnetsnmpxorp.so
%%XORP_SNMP%%mibs/ospf_mib_1850.la
%%XORP_SNMP%%mibs/ospf_mib_1850.so
%%XORP_SNMP%%mibs/xorp_if_mib_module.la
%%XORP_SNMP%%mibs/xorp_if_mib_module.so
%%XORP_SNMP%%mibs/snmpdscripts/getbgppeertable
%%XORP_SNMP%%mibs/snmpdscripts/getbgpversion
%%XORP_SNMP%%mibs/snmpdscripts/startsnmp
%%XORP_SNMP%%mibs/snmpdscripts/stopsnmp
@dirrm xrl/targets
@dirrm xrl
@dirrm static_routes
@dirrm rip/tools
@dirrm rip
@dirrm rib/tools
@dirrm rib
@dirrm policy
@dirrm pim
@dirrm ospf/tools
@dirrm ospf
@dirrm mld6igmp
%%XORP_SNMP%%@dirrm mibs/snmpdscripts
%%XORP_SNMP%%@dirrm mibs
@dirrm libxipc
%%XORP_SNMP%%@dirrmtry libdata/ldconfig/xorp
%%XORP_SNMP%%@dirrmtry libdata/ldconfig
%%XORP_SNMP%%@dirrmtry libdata
@dirrm fib2mrib
@dirrmtry fea/tools
@dirrm fea
@dirrmtry etc/templates
@dirrmtry etc/rc.d
@dirrmtry etc
@dirrm contrib/olsr/tools
@dirrm contrib/olsr
@dirrm contrib
@dirrm cli/tools
@dirrm cli
@dirrmtry bin
@dirrm bgp/tools
@dirrm bgp
