# mail/rmilter/pkg-plist
etc/rc.d/rmilter
@sample etc/rmilter.conf.sample
sbin/rmilter
@owner _rmilter
@group mail
@dir /var/run/rmilter
@owner
@group
