# www/oops/pkg-plist
sbin/oopsctl
sbin/oops
etc/oops/oops.cfg.sample
etc/oops/err_template.html.sample
etc/oops/auth_template.html.sample
etc/oops/passwd.sample
etc/oops/redir_rules.sample
etc/oops/redir_template.html.sample
etc/oops/accel_maps.sample
etc/oops/acl_local_networks.sample
etc/oops/select.sql.sample
etc/oops/tables/koi-alt.tab
etc/oops/tables/koi-iso.tab
etc/oops/tables/koi-win.tab
etc/rc.d/oops.sh.sample
libexec/oops/accel.so
libexec/oops/berkeley_db_api.so
libexec/oops/custom_log.so
libexec/oops/err.so
libexec/oops/fastredir.so
libexec/oops/gigabase_db_api.so
libexec/oops/lang.so
libexec/oops/log_dummy.so
libexec/oops/oopsctl_mod.so
libexec/oops/pam.so
libexec/oops/passwd_file.so
libexec/oops/passwd_mysql.so
libexec/oops/passwd_pgsql.so
libexec/oops/redir.so
libexec/oops/transparent.so
libexec/oops/vary_header.so
libexec/oops/wccp2.so
@dirrm libexec/oops
@dirrm etc/oops/tables
@dirrm etc/oops
@unexec rmdir %D/oops/DB 2>/dev/null || true
@unexec rmdir %D/oops/logs 2>/dev/null || true
@unexec rmdir %D/oops/storages 2>/dev/null || true
@unexec rmdir %D/oops 2>/dev/null || true
