# misc/raspbsd-led/files/openrc-flash-heartbeat-led.in
#!/sbin/openrc-run

name=flash_heartbeat_led
command="%%PREFIX%%/sbin/flash-heartbeat-led"

pidfile="/var/run/flash-heartbeat-led.pid"

depend()
{
	keyword -shutdown
}
