# sysutils/nwclient/pkg-plist
etc/rc.d/networker.sh
nsr/bin/backint
nsr/bin/init.utl
nsr/bin/lsm_backint
nsr/bin/mminfo
nsr/bin/mmlocate
nsr/bin/mmpool
nsr/bin/mmrecov
nsr/bin/nsradmin
nsr/bin/nsralist
nsr/bin/nsrarchive
nsr/bin/nsrclone
nsr/bin/nsrexec
nsr/bin/nsrexecd
nsr/bin/nsrfile
nsr/bin/nsrifmx
nsr/bin/nsrmm
nsr/bin/nsrports
nsr/bin/nsrretrieve
nsr/bin/nsrssc
nsr/bin/nsrstage
nsr/bin/nsrwatch
nsr/bin/perfmon.h
nsr/bin/preclntsave
nsr/bin/pstclntsave
nsr/bin/recover
nsr/bin/sapclone
nsr/bin/save
nsr/bin/savefs
nsr/bin/savegrp
nsr/bin/savepnpc
nsr/bin/uasm
nsr/gls/cm/registry
man/man8/EMASS_silo.8
man/man8/IBM_silo.8
man/man8/STK_silo.8
man/man8/ansrd.8
man/man8/changers.8
man/man8/dasadmin.8
man/man8/dbmi_config.8
man/man8/erase.8
man/man8/hpflip.8
man/man8/ielem.8
man/man8/inquire.8
man/man8/jb_config.8
man/man8/jbexercise.8
man/man8/jbinfo.8
man/man8/ldunld.8
man/man8/libscsi.8
man/man8/libsji.8
man/man8/libstlemass.8
man/man8/libstlibm.8
man/man8/libstlstk.8
man/man8/lrescan.8
man/man8/lreset.8
man/man8/lusbinfo.8
man/man8/lusdebug.8
man/man8/lusmode.8
man/man8/mini_el.8
man/man8/mm_data.5
man/man8/mminfo.8
man/man8/mmlocate.8
man/man8/mmpool.8
man/man8/mmrecov.8
man/man8/msense.8
man/man8/networker.8
man/man8/nsr.5
man/man8/nsr.8
man/man8/nsr_archive_request.5
man/man8/nsr_client.5
man/man8/nsr_crash.8
man/man8/nsr_data.5
man/man8/nsr_device.5
man/man8/nsr_directive.5
man/man8/nsr_getdate.3
man/man8/nsr_group.5
man/man8/nsr_ize.8
man/man8/nsr_jukebox.5
man/man8/nsr_label.5
man/man8/nsr_layout.5
man/man8/nsr_license.5
man/man8/nsr_migration.5
man/man8/nsr_notification.5
man/man8/nsr_policy.5
man/man8/nsr_pool.5
man/man8/nsr_regexp.5
man/man8/nsr_resource.5
man/man8/nsr_schedule.5
man/man8/nsr_service.5
man/man8/nsr_shutdown.8
man/man8/nsr_stage.5
man/man8/nsr_storage_node.5
man/man8/nsradmin.8
man/man8/nsralist.8
man/man8/nsrarchive.8
man/man8/nsrcap.8
man/man8/nsrcat.8
man/man8/nsrck.8
man/man8/nsrclone.8
man/man8/nsrd.8
man/man8/nsrdbmi.8
man/man8/nsrdmo.8
man/man8/nsrexec.8
man/man8/nsrexecd.8
man/man8/nsrib.8
man/man8/nsriba.8
man/man8/nsrim.8
man/man8/nsrindexasm.8
man/man8/nsrindexd.8
man/man8/nsrinfo.8
man/man8/nsrjb.8
man/man8/nsrlic.8
man/man8/nsrls.8
man/man8/nsrmm.8
man/man8/nsrmmd.8
man/man8/nsrmmdbasm.8
man/man8/nsrmmdbd.8
man/man8/nsrmon.8
man/man8/nsrnotd.8
man/man8/nsrports.8
man/man8/nsrretrieve.8
man/man8/nsrssc.8
man/man8/nsrstage.8
man/man8/nsrsyb.8
man/man8/nsrsybcc.8
man/man8/nsrsybrc.8
man/man8/nsrsybsv.8
man/man8/nsrtrap.8
man/man8/nsrwatch.8
man/man8/nwadmin.8
man/man8/nwarchive.8
man/man8/nwbackup.8
man/man8/nwrecover.8
man/man8/nwretrieve.8
man/man8/oraemcasm.8
man/man8/oraemcmap.8
man/man8/pmode.8
man/man8/preclntsave.8
man/man8/pstclntsave.8
man/man8/rap.8
man/man8/rapcheck.8
man/man8/rapd.8
man/man8/rapgen.1
man/man8/rapxfer.8
man/man8/recover.8
man/man8/relem.8
man/man8/resource.5
man/man8/save.8
man/man8/savefs.8
man/man8/savegrp.8
man/man8/savepnpc.8
man/man8/scanner.8
man/man8/sjidopen.8
man/man8/sjiielm.8
man/man8/sjiinq.8
man/man8/sjirdp.8
man/man8/sjirdtag.8
man/man8/sjirelem.8
man/man8/sjirjc.8
man/man8/ssi.8
man/man8/stk_eject.8
man/man8/stli.8
man/man8/tapeexercise.8
man/man8/tur.8
man/man8/uasm.8
man/man8/writebuf.8
nsr/gls/lc/os/portable/C
@dirrm nsr/bin
@dirrm nsr/gls/cm
@dirrm nsr/gls/lc/os/portable
@dirrm nsr/gls/lc/os
@dirrm nsr/gls/lc
@dirrm nsr/gls
@dirrm nsr
