# sysutils/nwclient/pkg-plist
etc/rc.d/networker.sh
nsr/bin/backint
nsr/bin/init.utl
nsr/bin/lsm_backint
nsr/bin/mminfo
nsr/bin/mmlocate
nsr/bin/mmpool
nsr/bin/mmrecov
nsr/bin/nsradmin
nsr/bin/nsralist
nsr/bin/nsrarchive
nsr/bin/nsrclone
nsr/bin/nsrexec
nsr/bin/nsrexecd
nsr/bin/nsrfile
nsr/bin/nsrifmx
nsr/bin/nsrmm
nsr/bin/nsrports
nsr/bin/nsrretrieve
nsr/bin/nsrssc
nsr/bin/nsrstage
nsr/bin/nsrwatch
nsr/bin/perfmon.h
nsr/bin/preclntsave
nsr/bin/pstclntsave
nsr/bin/recover
nsr/bin/sapclone
nsr/bin/save
nsr/bin/savefs
nsr/bin/savegrp
nsr/bin/savepnpc
nsr/bin/uasm
nsr/gls/cm/registry
nsr/gls/lc/os/portable/C
@dirrm nsr/bin
@dirrm nsr/gls/cm
@dirrm nsr/gls/lc/os/portable
@dirrm nsr/gls/lc/os
@dirrm nsr/gls/lc
@dirrm nsr/gls
@dirrm nsr
