# sysutils/bvm/pkg-plist
bin/bvm
bin/bvmb
@sample %%ETCDIR%%/bvm.conf.sample
@sample %%ETCDIR%%/nat.conf.sample
@sample %%ETCDIR%%/switch.conf.sample
etc/rc.d/bvmd
