# www/pserv/pkg-plist
sbin/pserv
etc/mime.types
etc/pserv.conf
etc/rc.d/pserv.sh
