# www/pserv/pkg-plist
sbin/pserv
etc/pserv/mime.types
etc/pserv/pserv.conf
etc/rc.d/pserv.sh
@dirrm etc/pserv
