# misc/spacenavd-devel/pkg-plist
bin/spacenavd
bin/spnavd_ctl
etc/rc.d/spacenavd
@sample etc/spnavrc.sample
