# net/smokeping/pkg-plist
bin/smokeping
etc/rc.d/smokeping.sh
@unexec if cmp -s %D/etc/smokeping/basepage.html %D/etc/smokeping/basepage.html.dist; then rm -f %D/etc/smokeping/basepage.html; fi
etc/smokeping/basepage.html.dist
@exec [ -f %B/basepage.html ] || cp %B/%f %B/basepage.html
@unexec if cmp -s %D/etc/smokeping/config %D/etc/smokeping/config.dist; then rm -f %D/etc/smokeping/config; fi
etc/smokeping/config.dist
@exec [ -f %B/config ] || cp %B/%f %B/config
@unexec if cmp -s %D/etc/smokeping/smokemail %D/etc/smokeping/smokemail.dist; then rm -f %D/etc/smokeping/smokemail; fi
etc/smokeping/smokemail.dist
@exec [ -f %B/smokemail ] || cp %B/%f %B/smokemail
@unexec if cmp -s %D/etc/smokeping/config-echoping %D/etc/smokeping/config-echoping.dist; then rm -f %D/etc/smokeping/config-echoping; fi
etc/smokeping/config-echoping.dist
@exec [ -f %B/config-echoping ] || cp %B/%f %B/config-echoping
%%PORTDOCS%%%%DOCSDIR%%/CHANGES
%%PORTDOCS%%%%DOCSDIR%%/CONTRIBUTORS
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%%%DOCSDIR%%/CiscoRTTMonDNS.pm.html
%%PORTDOCS%%%%DOCSDIR%%/CiscoRTTMonDNS.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/CiscoRTTMonEchoICMP.pm.html
%%PORTDOCS%%%%DOCSDIR%%/CiscoRTTMonEchoICMP.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/CiscoRTTMonTcpConnect.pm.html
%%PORTDOCS%%%%DOCSDIR%%/CiscoRTTMonTcpConnect.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/Curl.pm.html
%%PORTDOCS%%%%DOCSDIR%%/Curl.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/DNS.pm.html
%%PORTDOCS%%%%DOCSDIR%%/DNS.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/EchoPing.pm.html
%%PORTDOCS%%%%DOCSDIR%%/EchoPing.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/EchoPingChargen.pm.html
%%PORTDOCS%%%%DOCSDIR%%/EchoPingChargen.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/EchoPingDiscard.pm.html
%%PORTDOCS%%%%DOCSDIR%%/EchoPingDiscard.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/EchoPingHttp.pm.html
%%PORTDOCS%%%%DOCSDIR%%/EchoPingHttp.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/EchoPingHttps.pm.html
%%PORTDOCS%%%%DOCSDIR%%/EchoPingHttps.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/EchoPingIcp.pm.html
%%PORTDOCS%%%%DOCSDIR%%/EchoPingIcp.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/EchoPingSmtp.pm.html
%%PORTDOCS%%%%DOCSDIR%%/EchoPingSmtp.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/FPing.pm.html
%%PORTDOCS%%%%DOCSDIR%%/FPing.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/FPing6.pm.html
%%PORTDOCS%%%%DOCSDIR%%/FPing6.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/IOSPing.pm.html
%%PORTDOCS%%%%DOCSDIR%%/IOSPing.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/ParseConfig.pm.html
%%PORTDOCS%%%%DOCSDIR%%/ParseConfig.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/RemoteFPing.pm.html
%%PORTDOCS%%%%DOCSDIR%%/RemoteFPing.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping.pm.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/base.pm.html
%%PORTDOCS%%%%DOCSDIR%%/base.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/basefork.pm.html
%%PORTDOCS%%%%DOCSDIR%%/basefork.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/basevars.pm.html
%%PORTDOCS%%%%DOCSDIR%%/basevars.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping.cgi.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping.cgi.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping_config.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping_config.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping_install.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping_install.txt
%%PORTDOCS%%%%DOCSDIR%%/telnetIOSPing.pm.html
%%PORTDOCS%%%%DOCSDIR%%/telnetIOSPing.pm.txt
smokeping/htdocs/smokeping.cgi
smokeping/lib/ISG/ParseConfig.pm
smokeping/lib/Smokeping.pm
smokeping/lib/ciscoRttMonMIB.pm
smokeping/lib/probes/AnotherDNS.pm
smokeping/lib/probes/AnotherSSH.pm
smokeping/lib/probes/CiscoRTTMonDNS.pm
smokeping/lib/probes/CiscoRTTMonEchoICMP.pm
smokeping/lib/probes/CiscoRTTMonTcpConnect.pm
smokeping/lib/probes/Curl.pm
smokeping/lib/probes/DNS.pm
smokeping/lib/probes/EchoPing.pm
smokeping/lib/probes/EchoPingChargen.pm
smokeping/lib/probes/EchoPingDiscard.pm
smokeping/lib/probes/EchoPingHttp.pm
smokeping/lib/probes/EchoPingHttps.pm
smokeping/lib/probes/EchoPingIcp.pm
smokeping/lib/probes/EchoPingSmtp.pm
smokeping/lib/probes/FPing.pm
smokeping/lib/probes/FPing6.pm
smokeping/lib/probes/IOSPing.pm
smokeping/lib/probes/LDAP.pm
smokeping/lib/probes/Radius.pm
smokeping/lib/probes/RemoteFPing.pm
smokeping/lib/probes/SSH.pm
smokeping/lib/probes/base.pm
smokeping/lib/probes/basefork.pm
smokeping/lib/probes/basevars.pm
smokeping/lib/probes/passwordchecker.pm
smokeping/lib/probes/telnetIOSPing.pm
@unexec rmdir %D/etc/smokeping 2>/dev/null || true
@exec mkdir -p %D/smokeping/htdocs/img; chown www:www %D/smokeping/htdocs/img
@unexec rm -rf %D/smokeping/htdocs/img 2>/dev/null || true
@exec mkdir -p %D/var/smokeping
@exec chown smokeping:smokeping %D/var/smokeping
@unexec rmdir %D/var/smokeping 2>/dev/null || true
@unexec rmdir %D/var 2>/dev/null || true
@dirrm smokeping/htdocs
@dirrm smokeping/lib/probes
@dirrm smokeping/lib/ISG
@dirrm smokeping/lib
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@unexec rmdir %D/smokeping 2>/dev/null || true
