# net/smokeping/pkg-plist
bin/smokeping
etc/rc.d/smokeping.sh
@unexec if cmp -s %D/etc/smokeping/basepage.html %D/etc/smokeping/basepage.html.dist; then rm -f %D/etc/smokeping/basepage.html; fi
etc/smokeping/basepage.html.dist
@exec [ -f %B/basepage.html ] || cp %B/%f %B/basepage.html
@unexec if cmp -s %D/etc/smokeping/config %D/etc/smokeping/config.dist; then rm -f %D/etc/smokeping/config; fi
etc/smokeping/config.dist
@exec [ -f %B/config ] || cp %B/%f %B/config
@unexec if cmp -s %D/etc/smokeping/smokemail %D/etc/smokeping/smokemail.dist; then rm -f %D/etc/smokeping/smokemail; fi
etc/smokeping/smokemail.dist
@exec [ -f %B/smokemail ] || cp %B/%f %B/smokemail
@unexec if cmp -s %D/etc/smokeping/config-echoping %D/etc/smokeping/config-echoping.dist; then rm -f %D/etc/smokeping/config-echoping; fi
etc/smokeping/config-echoping.dist
@exec [ -f %B/config-echoping ] || cp %B/%f %B/config-echoping
%%PORTDOCS%%%%DOCSDIR%%/CHANGES
%%PORTDOCS%%%%DOCSDIR%%/CONTRIBUTORS
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%%%DOCSDIR%%/ParseConfig.pm.html
%%PORTDOCS%%%%DOCSDIR%%/ParseConfig.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping.pm.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping.cgi.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping.cgi.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping_config.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping_config.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping_install.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping_install.txt
%%PORTDOCS%%%%DOCSDIR%%/matchers/avgratio.pm.html
%%PORTDOCS%%%%DOCSDIR%%/matchers/avgratio.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/matchers/base.pm.html
%%PORTDOCS%%%%DOCSDIR%%/matchers/base.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/matchers/median.pm.html
%%PORTDOCS%%%%DOCSDIR%%/matchers/median.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/AnotherDNS.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/AnotherDNS.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/AnotherSSH.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/AnotherSSH.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/CiscoRTTMonDNS.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/CiscoRTTMonDNS.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/CiscoRTTMonEchoICMP.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/CiscoRTTMonEchoICMP.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/CiscoRTTMonTcpConnect.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/CiscoRTTMonTcpConnect.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/Curl.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/Curl.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/DNS.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/DNS.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPing.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPing.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPingChargen.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPingChargen.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPingDiscard.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPingDiscard.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPingHttp.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPingHttp.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPingHttps.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPingHttps.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPingIcp.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPingIcp.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPingSmtp.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/EchoPingSmtp.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/FPing.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/FPing.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/FPing6.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/FPing6.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/IOSPing.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/IOSPing.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/LDAP.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/LDAP.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/Radius.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/Radius.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/RemoteFPing.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/RemoteFPing.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/SSH.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/SSH.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/base.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/base.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/basefork.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/basefork.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/basevars.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/basevars.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/passwordchecker.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/passwordchecker.pm.txt
%%PORTDOCS%%%%DOCSDIR%%/probes/telnetIOSPing.pm.html
%%PORTDOCS%%%%DOCSDIR%%/probes/telnetIOSPing.pm.txt
smokeping/htdocs/smokeping.cgi
smokeping/lib/ISG/ParseConfig.pm
smokeping/lib/Smokeping.pm
smokeping/lib/ciscoRttMonMIB.pm
smokeping/lib/matchers/base.pm
smokeping/lib/matchers/median.pm
smokeping/lib/matchers/avgratio.pm
smokeping/lib/probes/AnotherDNS.pm
smokeping/lib/probes/AnotherSSH.pm
smokeping/lib/probes/CiscoRTTMonDNS.pm
smokeping/lib/probes/CiscoRTTMonEchoICMP.pm
smokeping/lib/probes/CiscoRTTMonTcpConnect.pm
smokeping/lib/probes/Curl.pm
smokeping/lib/probes/DNS.pm
smokeping/lib/probes/EchoPing.pm
smokeping/lib/probes/EchoPingChargen.pm
smokeping/lib/probes/EchoPingDiscard.pm
smokeping/lib/probes/EchoPingHttp.pm
smokeping/lib/probes/EchoPingHttps.pm
smokeping/lib/probes/EchoPingIcp.pm
smokeping/lib/probes/EchoPingSmtp.pm
smokeping/lib/probes/FPing.pm
smokeping/lib/probes/FPing6.pm
smokeping/lib/probes/IOSPing.pm
smokeping/lib/probes/LDAP.pm
smokeping/lib/probes/Radius.pm
smokeping/lib/probes/RemoteFPing.pm
smokeping/lib/probes/SSH.pm
smokeping/lib/probes/base.pm
smokeping/lib/probes/basefork.pm
smokeping/lib/probes/basevars.pm
smokeping/lib/probes/passwordchecker.pm
smokeping/lib/probes/telnetIOSPing.pm
@unexec rmdir %D/etc/smokeping 2>/dev/null || true
@exec mkdir -p %D/smokeping/htdocs/img; chown www:www %D/smokeping/htdocs/img
@unexec rm -rf %D/smokeping/htdocs/img 2>/dev/null || true
@exec mkdir -p %D/var/smokeping
@exec chown smokeping:smokeping %D/var/smokeping
@unexec rmdir %D/var/smokeping 2>/dev/null || true
@unexec rmdir %D/var 2>/dev/null || true
@dirrm smokeping/htdocs
@dirrm smokeping/lib/probes
@dirrm smokeping/lib/matchers
@dirrm smokeping/lib/ISG
@dirrm smokeping/lib
%%PORTDOCS%%@dirrm %%DOCSDIR%%/probes
%%PORTDOCS%%@dirrm %%DOCSDIR%%/matchers
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@unexec rmdir %D/smokeping 2>/dev/null || true
