# net-mgmt/fprobe/files/fprobe.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: fprobe
# REQUIRE: NETWORKING
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable %%NAME%%:
#
#%%NAME%%_enable="YES"
#
# See fprobe(8) for %%NAME%%_flags
#

. /etc/rc.subr

name="%%NAME%%"
rcvar=%%NAME%%_enable

load_rc_config "$name"

: ${fprobe_enable="NO"}
: ${fprobe_flags:="-i lo0 127.0.0.1:2055"}

command="%%PREFIX%%/sbin/fprobe"
pidfile="/var/run/$name.pid"

run_rc_command "$1"
