# security/clamav/pkg-plist
bin/clamav-config
bin/clambc
bin/clamconf
bin/clamdscan
bin/clamdtop
bin/clamscan
bin/clamsubmit
bin/freshclam
bin/sigtool
etc/certs/clamav.crt
%%MILTER%%@sample etc/clamav-milter.conf.sample
@sample etc/clamd.conf.sample
@sample etc/freshclam.conf.sample
%%MILTER%%etc/rc.d/clamav_milter
include/clamav-config.h
include/clamav-types.h
include/clamav-version.h
include/clamav.h
include/libfreshclam.h
lib/libclamav.so
lib/libclamav.so.12
lib/libclamav.so.12.1.0
lib/libclamav_rust.a
lib/libclamunrar.so
lib/libclamunrar.so.12
lib/libclamunrar.so.12.1.0
lib/libclamunrar_iface.so
lib/libclamunrar_iface.so.12
lib/libclamunrar_iface.so.12.1.0
lib/libfreshclam.so
lib/libfreshclam.so.4
lib/libfreshclam.so.4.0.0
libdata/pkgconfig/libclamav.pc
%%MILTER%%sbin/clamav-milter
sbin/clamd
%%PORTDOCS%%%%DOCSDIR%%/html/404.html
%%PORTDOCS%%%%DOCSDIR%%/html/Introduction.html
%%PORTDOCS%%%%DOCSDIR%%/html/ace-2a3cd908.js
%%PORTDOCS%%%%DOCSDIR%%/html/appendix/Appendix.html
%%PORTDOCS%%%%DOCSDIR%%/html/appendix/Authenticode.html
%%PORTDOCS%%%%DOCSDIR%%/html/appendix/CvdPrivateMirror.html
%%PORTDOCS%%%%DOCSDIR%%/html/appendix/FileTypes.html
%%PORTDOCS%%%%DOCSDIR%%/html/appendix/FunctionalityLevels.html
%%PORTDOCS%%%%DOCSDIR%%/html/appendix/Terminology.html
%%PORTDOCS%%%%DOCSDIR%%/html/ayu-highlight-3fdfc3ac.css
%%PORTDOCS%%%%DOCSDIR%%/html/book-a0b12cfe.js
%%PORTDOCS%%%%DOCSDIR%%/html/clipboard-1626706a.min.js
%%PORTDOCS%%%%DOCSDIR%%/html/community_resources/CommunityResources.html
%%PORTDOCS%%%%DOCSDIR%%/html/community_resources/CompileClamAV_AmazonLinux2.txt
%%PORTDOCS%%%%DOCSDIR%%/html/css/chrome-ae938929.css
%%PORTDOCS%%%%DOCSDIR%%/html/css/general-8eed8816.css
%%PORTDOCS%%%%DOCSDIR%%/html/css/print-9e4910d8.css
%%PORTDOCS%%%%DOCSDIR%%/html/css/variables-6d03f665.css
%%PORTDOCS%%%%DOCSDIR%%/html/editor-16ca416c.js
%%PORTDOCS%%%%DOCSDIR%%/html/elasticlunr-ef4e11c1.min.js
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-cvd.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-eol.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-freshclam.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-ignore.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-malware-fp-reports.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-misc.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-ml.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-pua.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-rust.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-safebrowsing.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-scan-alerts.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-troubleshoot.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-uninstall.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-upgrade.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-whichversion.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq-win32.html
%%PORTDOCS%%%%DOCSDIR%%/html/faq/faq.html
%%PORTDOCS%%%%DOCSDIR%%/html/favicon-ba1b3133.png
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/OPEN-SANS-LICENSE-cfc7749b.txt
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/SOURCE-CODE-PRO-LICENSE-d1e6d465.txt
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/fonts-9644e21d.css
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/open-sans-v17-all-charsets-300-7736aa35.woff2
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/open-sans-v17-all-charsets-300italic-2c7b95c0.woff2
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/open-sans-v17-all-charsets-600-486c6759.woff2
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/open-sans-v17-all-charsets-600italic-1a3e8659.woff2
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/open-sans-v17-all-charsets-700-c22fe8c7.woff2
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/open-sans-v17-all-charsets-700italic-238ae959.woff2
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/open-sans-v17-all-charsets-800-3d2c812a.woff2
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/open-sans-v17-all-charsets-800italic-ba1521ec.woff2
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/open-sans-v17-all-charsets-italic-6c9463f7.woff2
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/open-sans-v17-all-charsets-regular-2e3b1d34.woff2
%%PORTDOCS%%%%DOCSDIR%%/html/fonts/source-code-pro-v11-all-charsets-500-2bdd9410.woff2
%%PORTDOCS%%%%DOCSDIR%%/html/googled62299e9391332c4.html
%%PORTDOCS%%%%DOCSDIR%%/html/highlight-493f70e1.css
%%PORTDOCS%%%%DOCSDIR%%/html/highlight-abc7f01d.js
%%PORTDOCS%%%%DOCSDIR%%/html/images/change-fork-name.png
%%PORTDOCS%%%%DOCSDIR%%/html/images/cisco.png
%%PORTDOCS%%%%DOCSDIR%%/html/images/clamav-git-workflow.png
%%PORTDOCS%%%%DOCSDIR%%/html/images/clone-your-fork.png
%%PORTDOCS%%%%DOCSDIR%%/html/images/create-a-fork.png
%%PORTDOCS%%%%DOCSDIR%%/html/images/demon.png
%%PORTDOCS%%%%DOCSDIR%%/html/images/flamegraph.svg
%%PORTDOCS%%%%DOCSDIR%%/html/images/fork-is-behind.png
%%PORTDOCS%%%%DOCSDIR%%/html/images/logo.png
%%PORTDOCS%%%%DOCSDIR%%/html/index.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Development.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Development/Contribute.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Development/build-installer-packages.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Development/clamav-git-work-flow.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Development/code-coverage.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Development/development-builds.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Development/fuzzing-sanitizers.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Development/github-pr-basics.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Development/libclamav.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Development/performance-profiling.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Development/personal-forks.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Development/testing-pull-requests.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Development/tips-and-tricks.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Installing.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Installing/Add-clamav-user.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Installing/Community-projects.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Installing/Docker.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Installing/Installing-from-source-Unix.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Installing/Installing-from-source-Windows.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Installing/Packages.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/OnAccess.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/AllowLists.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/AuthenticodeRules.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/BodySignatureFormat.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/BytecodeSignatures.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/ContainerMetadata.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/DatabaseInfo.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/DynamicConfig.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/EncryptedArchives.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/ExtendedSignatures.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/FileTypeMagic.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/HashSignatures.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/LogicalSignatures.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/PhishSigs.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/SignatureNames.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Signatures/YaraRules.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Usage.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Usage/Configuration.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Usage/ReportABug.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Usage/Scanning.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Usage/Services.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/Usage/SignatureManagement.html
%%PORTDOCS%%%%DOCSDIR%%/html/manual/cisco-talos.gpg
%%PORTDOCS%%%%DOCSDIR%%/html/mark-09e88c2c.min.js
%%PORTDOCS%%%%DOCSDIR%%/html/mode-rust-2c9d5c9a.js
%%PORTDOCS%%%%DOCSDIR%%/html/print.html
%%PORTDOCS%%%%DOCSDIR%%/html/searcher-c2a407aa.js
%%PORTDOCS%%%%DOCSDIR%%/html/searchindex-1b5ba28b.js
%%PORTDOCS%%%%DOCSDIR%%/html/sitemap.xml
%%PORTDOCS%%%%DOCSDIR%%/html/theme-dawn-4493f9c8.js
%%PORTDOCS%%%%DOCSDIR%%/html/theme-tomorrow_night-9dbe62a9.js
%%PORTDOCS%%%%DOCSDIR%%/html/toc-cbaddea7.js
%%PORTDOCS%%%%DOCSDIR%%/html/toc.html
%%PORTDOCS%%%%DOCSDIR%%/html/tomorrow-night-4c0ae647.css
share/man/man1/clambc.1.gz
share/man/man1/clamconf.1.gz
share/man/man1/clamdscan.1.gz
share/man/man1/clamdtop.1.gz
share/man/man1/clamscan.1.gz
share/man/man1/clamsubmit.1.gz
share/man/man1/freshclam.1.gz
share/man/man1/sigtool.1.gz
share/man/man5/clamav-milter.conf.5.gz
share/man/man5/clamd.conf.5.gz
share/man/man5/freshclam.conf.5.gz
share/man/man8/clamav-milter.8.gz
share/man/man8/clamd.8.gz
@dir(%%CLAMAVUSER%%,%%CLAMAVGROUP%%,0755) %%DBDIR%%
@dir(%%CLAMAVUSER%%,%%CLAMAVGROUP%%,0755) %%LOGDIR%%
@dir(%%CLAMAVUSER%%,%%CLAMAVGROUP%%,0755) %%RUNDIR%%
