# security/clamav/pkg-plist
@comment $FreeBSD$
@unexec [ ! -f /var/run/clamav/clamd.pid ] || %%RC_DIR%%/clamd%%RC_SUFX%% stop || true
bin/clamscan
bin/clamdscan
bin/freshclam
bin/sigtool
%%CLAMAV-MILTER%%etc/rc.d/clamav-milter%%RC_SUFX%%
%%CLAMAV-MILTER%%sbin/clamav-milter
sbin/clamd
include/clamav.h
lib/libclamav.so.1
lib/libclamav.so
lib/libclamav.a
@unexec if cmp -s %D/etc/clamav.conf %D/etc/clamav.conf.default; then rm -f %D/etc/clamav.conf; fi
etc/clamav.conf.default
@exec [ -f %B/clamav.conf ] || cp %B/%f %B/clamav.conf
etc/rc.d/clamav-clamd%%RC_SUFX%%
etc/rc.d/clamav-freshclam%%RC_SUFX%%
@owner clamav
@group clamav
%%DATADIR%%/daily.cvd
%%DATADIR%%/main.cvd
%%DATADIR%%/mirrors.txt
@exec chown clamav:clamav %B
@dirrm %%DATADIR%%
@unexec rmdir /var/run/clamav 2>/dev/null || true
@unexec rmdir /var/log/clamav 2>/dev/null || true
