# www/mini_httpd/pkg-plist
@sample etc/mini_httpd.cnf.sample
etc/rc.d/mini_httpd
man/man1/mini_httpd-htpasswd.1.gz
man/man8/mini_httpd.8.gz
sbin/mini_httpd
sbin/mini_httpd-htpasswd
sbin/mini_httpd_wrapper
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/500.mini_httpd-rotate
%%PORTEXAMPLES%%@dir %%EXAMPLESDIR%%
