# net/gnu-finger/pkg-plist
bin/finger
libexec/fingerd
libexec/in.fingerd
libexec/in.cfingerd
etc/rc.d/fingerd.sh.sample
etc/fingerdir/ttylocs.sample
etc/fingerdir/clients.sample
etc/fingerdir/mailhost.sample
etc/fingerdir/serverhost.sample
etc/fingerdir/targets/x-.help
etc/fingerdir/targets/x-.site
@dirrm etc/fingerdir/targets
@dirrm etc/fingerdir
