# audio/mpdscribble/files/openrc-mpdscribble.in
#!/sbin/openrc-run

name=mpdscribble
command=%%PREFIX%%/bin/mpdscribble

: ${mpdscribble_user="nobody"}

command_args="--daemon-user ${mpdscribble_user}"
command_user=${mpdscribble_user}

depend() {
	keyword -shutdown
}
