# sysutils/iocell/pkg-plist
etc/rc.d/iocell
lib/iocell/ioc-cmd
lib/iocell/ioc-common
lib/iocell/ioc-configure
lib/iocell/ioc-git
lib/iocell/ioc-globals
lib/iocell/ioc-hacks
lib/iocell/ioc-help
lib/iocell/ioc-image
lib/iocell/ioc-info
lib/iocell/ioc-network
lib/iocell/ioc-rc
lib/iocell/ioc-rctl
lib/iocell/ioc-send
lib/iocell/ioc-snapshot
lib/iocell/ioc-upgrade
lib/iocell/ioc-zfs
man/man8/iocell.8.gz
sbin/iocell
