# mail/opendmarc/files/opendmarc.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: opendmarc
# REQUIRE: LOGIN
# BEFORE:  mail
# KEYWORD: opendmarc

# Define these opendmarc_* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#       /etc/rc.conf.d/opendmarc
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#

. /etc/rc.subr

name="opendmarc"
rcvar=opendmarc_enable

command="%%PREFIX%%/sbin/opendmarc"

load_rc_config $name

opendmarc_enable=${opendmarc_enable-"NO"}
opendmarc_runas=${opendmarc_runas-"mailnull:mailnull"}
opendmarc_pidfile=${opendmarc_pidfile-"/var/run/opendmarc.pid"}
opendmarc_socketspec=${opendmarc_socketspec-"inet:8893@localhost"}
opendmarc_cfgfile=${opendmarc_cfgfile-"%%PREFIX%%/etc/mail/opendmarc.conf"}
opendmarc_flags=${opendmarc_flags-"-l -P $opendmarc_pidfile \
 -c $opendmarc_cfgfile -p $opendmarc_socketspec -u $opendmarc_runas"}

run_rc_command "$1"
