# www/mnogosearch31/pkg-plist
etc/mnogosearch/indexer.conf-dist
etc/mnogosearch/search.htm-dist
etc/mnogosearch/spelld.conf-dist
etc/rc.d/spelld.sh-dist
include/udm_agent.h
include/udm_alias.h
include/udm_boolean.h
include/udm_cache.h
include/udm_charset.h
include/udm_common.h
include/udm_conf.h
include/udm_config.h
include/udm_crc32.h
include/udm_crossword.h
include/udm_db.h
include/udm_env.h
include/udm_execget.h
include/udm_filter.h
include/udm_ftp.h
include/udm_host.h
include/udm_hrefs.h
include/udm_http.h
include/udm_id3.h
include/udm_indexer.h
include/udm_log.h
include/udm_mimetype.h
include/udm_mutex.h
include/udm_parsedate.h
include/udm_parsehtml.h
include/udm_parser.h
include/udm_parseurl.h
include/udm_proto.h
include/udm_robots.h
include/udm_search_tl.h
include/udm_searchtool.h
include/udm_server.h
include/udm_services.h
include/udm_signals.h
include/udm_socket.h
include/udm_spell.h
include/udm_stopwords.h
include/udm_utils.h
include/udm_word.h
include/udm_xmalloc.h
include/udmsearch.h
lib/libudmsearch.a
lib/libudmsearch.la
lib/libudmsearch.so
lib/libudmsearch.so.1
sbin/cachelogd
sbin/indexer
sbin/run-splitter
sbin/spelld
sbin/splitter
bin/udm-config
share/doc/mnogosearch/HTDB.txt
share/doc/mnogosearch/HTML-design.txt
share/doc/mnogosearch/HTTP-codes.txt
share/doc/mnogosearch/INSTALL
share/doc/mnogosearch/NEWS.html
share/doc/mnogosearch/alias.txt
share/doc/mnogosearch/bugs.txt
share/doc/mnogosearch/cache.txt
share/doc/mnogosearch/cachemode.txt
share/doc/mnogosearch/categories.txt
share/doc/mnogosearch/charset.txt
share/doc/mnogosearch/follow.txt
share/doc/mnogosearch/hrefonly.txt
share/doc/mnogosearch/html.txt
share/doc/mnogosearch/ispell.txt
share/doc/mnogosearch/lib.txt
share/doc/mnogosearch/mnoGoSearch-Oracle8.html
share/doc/mnogosearch/mp3.txt
share/doc/mnogosearch/negotiation.txt
share/doc/mnogosearch/parsers.txt
share/doc/mnogosearch/performance.txt
share/doc/mnogosearch/relevancy.txt
share/doc/mnogosearch/search.txt
share/doc/mnogosearch/srv_table.txt
share/doc/mnogosearch/storage.txt
share/doc/mnogosearch/syslog.txt
share/doc/mnogosearch/tags.txt
share/doc/mnogosearch/templates.txt
share/doc/mnogosearch/tracking.txt
share/doc/mnogosearch/usage.txt
share/mnogosearch/create/ibase/create.sql
share/mnogosearch/create/msql/crc-multi.txt
share/mnogosearch/create/msql/crc.txt
share/mnogosearch/create/msql/create.txt
share/mnogosearch/create/msql/cross-crc.txt
share/mnogosearch/create/msql/cross.txt
share/mnogosearch/create/msql/ispell.txt
share/mnogosearch/create/msql/multi.txt
share/mnogosearch/create/msql/track.txt
share/mnogosearch/create/mssql/create.sql
share/mnogosearch/create/mysql/categories.txt
share/mnogosearch/create/mysql/crc-multi.txt
share/mnogosearch/create/mysql/crc.txt
share/mnogosearch/create/mysql/create.txt
share/mnogosearch/create/mysql/cross-crc.txt
share/mnogosearch/create/mysql/cross.txt
share/mnogosearch/create/mysql/ispell.txt
share/mnogosearch/create/mysql/multi.txt
share/mnogosearch/create/mysql/servers.txt
share/mnogosearch/create/mysql/track.txt
share/mnogosearch/create/mysql/url-raid.txt
share/mnogosearch/create/oracle/categories.sql
share/mnogosearch/create/oracle/crc-multi.sql
share/mnogosearch/create/oracle/crc.sql
share/mnogosearch/create/oracle/create.sql
share/mnogosearch/create/oracle/multi.sql
share/mnogosearch/create/oracle/servers.sql
share/mnogosearch/create/oracle/track.sql
share/mnogosearch/create/pgsql/categories.txt
share/mnogosearch/create/pgsql/crc-multi.txt
share/mnogosearch/create/pgsql/crc.txt
share/mnogosearch/create/pgsql/create.txt
share/mnogosearch/create/pgsql/cross-crc.txt
share/mnogosearch/create/pgsql/cross.txt
share/mnogosearch/create/pgsql/grant.txt
share/mnogosearch/create/pgsql/ispell.txt
share/mnogosearch/create/pgsql/multi.txt
share/mnogosearch/create/pgsql/news-extension.txt
share/mnogosearch/create/pgsql/server.txt
share/mnogosearch/create/pgsql/to_use_with_tags.txt
share/mnogosearch/create/pgsql/track.txt
share/mnogosearch/create/sapdb/README
share/mnogosearch/create/sapdb/categories.ins
share/mnogosearch/create/sapdb/crc-multi.ins
share/mnogosearch/create/sapdb/crc.ins
share/mnogosearch/create/sapdb/create-all.ins
share/mnogosearch/create/sapdb/create.ins
share/mnogosearch/create/sapdb/create.sh
share/mnogosearch/create/sapdb/drop.ins
share/mnogosearch/create/sapdb/multi.ins
share/mnogosearch/create/sapdb/track.ins
share/mnogosearch/create/solid/README
share/mnogosearch/create/solid/create.txt
share/mnogosearch/create/solid/multi.txt
share/mnogosearch/create/stopwords/stop.cz.txt
share/mnogosearch/create/stopwords/stop.de.txt
share/mnogosearch/create/stopwords/stop.dk.txt
share/mnogosearch/create/stopwords/stop.en.big
share/mnogosearch/create/stopwords/stop.en.huge
share/mnogosearch/create/stopwords/stop.en.txt
share/mnogosearch/create/stopwords/stop.es.txt
share/mnogosearch/create/stopwords/stop.fr.txt
share/mnogosearch/create/stopwords/stop.it.txt
share/mnogosearch/create/stopwords/stop.nl.txt
share/mnogosearch/create/stopwords/stop.no.txt
share/mnogosearch/create/stopwords/stop.pl.txt
share/mnogosearch/create/stopwords/stop.pt.txt
share/mnogosearch/create/stopwords/stop.ru.txt
share/mnogosearch/create/stopwords/stop.sk.txt
share/mnogosearch/create/stopwords/stop.tr.txt
share/mnogosearch/create/stopwords/stop.ua.txt
share/mnogosearch/create/stopwords/stopword.pl
share/mnogosearch/create/sybase/README
share/mnogosearch/create/sybase/create.sql
share/mnogosearch/create/virtuoso/crc.txt
share/mnogosearch/create/virtuoso/create.txt
share/mnogosearch/create/virtuoso/multi.txt
share/mnogosearch/search.cgi
@dirrm share/doc/mnogosearch
@dirrm share/mnogosearch/create/ibase
@dirrm share/mnogosearch/create/msql
@dirrm share/mnogosearch/create/mssql
@dirrm share/mnogosearch/create/mysql
@dirrm share/mnogosearch/create/oracle
@dirrm share/mnogosearch/create/pgsql
@dirrm share/mnogosearch/create/sapdb
@dirrm share/mnogosearch/create/solid
@dirrm share/mnogosearch/create/stopwords
@dirrm share/mnogosearch/create/sybase
@dirrm share/mnogosearch/create/virtuoso
@dirrm share/mnogosearch/create
@dirrm share/mnogosearch
@dirrm etc/mnogosearch
