# www/perlbal/files/perlbal.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: perlbal
# REQUIRE: LOGIN
# KEYWORK: shutdown

. /etc/rc.subr

name="perlbal"
command="%%PREFIX%%/bin/perlbal"

load_rc_config "$name"

: ${perlbal_enable="NO"}
: ${perlbal_configuration="%%PREFIX%%/etc/perlbal/perlbal.conf"}

required_files="$perlbal_configuration"

command_args="--config $required_files --daemon"

run_rc_command "$1"

