# news/dnews/pkg-plist
bin/tellnews
etc/rc.d/dnews.sh
lib/dnews/access.conf.sample
lib/dnews/control.conf.sample
lib/dnews/dnews
lib/dnews/dnews.conf.sample
lib/dnews/dnews_start
lib/dnews/doc/faq.html
lib/dnews/doc/install.txt
lib/dnews/doc/license.txt
lib/dnews/doc/manual.txt
lib/dnews/dposter
lib/dnews/expire.conf.sample
lib/dnews/fts.conf.sample
lib/dnews/hp.conf.sample
lib/dnews/moderators.conf.sample
lib/dnews/newsfeeds.conf.sample
lib/dnews/nocached.body
lib/dnews/nocached.head
lib/dnews/scripts/newver.sh
lib/dnews/xmit
