# news/dnews/pkg-plist
bin/tellnews
dnews/access.conf
dnews/control.conf
dnews/dfeed
dnews/dis.sh
dnews/distrib.sh
dnews/dmgr/access.htm
dnews/dmgr/adv.htm
dnews/dmgr/config.htm
dnews/dmgr/dir.htm
dnews/dmgr/dmgrcgi.cgi
dnews/dmgr/feed.htm
dnews/dmgr/feeds.htm
dnews/dmgr/help.htm
dnews/dmgr/log.htm
dnews/dmgr/logoff.htm
dnews/dmgr/menu.htm
dnews/dmgr/sucks.htm
dnews/dmgr/tellnews.htm
dnews/dmgr/users.htm
dnews/dmgrsvc
dnews/dmulti
dnews/dmulti_start.sh
dnews/dnews
dnews/dnews_start
dnews/dnews_start.sh
dnews/dnewsauth.sh
dnews/dnewssrt
dnews/dnewsweb
dnews/dnewsweb.ini
dnews/dnewsweb.txt
dnews/dposter
dnews/drobot
dnews/dudp
dnews/duse
dnews/expire.conf
dnews/faq.html
dnews/fixinet.sh
dnews/fixowner.sh
dnews/fts.conf
dnews/gupcgi
dnews/hp.conf
dnews/install.sh
dnews/install.txt
dnews/libauth.sh
dnews/license.txt
dnews/lint.sh
dnews/makelib.sh
dnews/manual.txt
dnews/mkgroups.sh
dnews/moderators.conf
dnews/newsfeeds.conf
dnews/newver.sh
dnews/nocached.body
dnews/resetxxx.sh
dnews/restart.sh
dnews/setfeed.ini
dnews/setprot.sh
dnews/tar.sh
dnews/tellnews
dnews/test.sh
dnews/testing.sh
dnews/web/advert.htm
dnews/web/cancel.tpl
dnews/web/config.tpl
dnews/web/empty.htm
dnews/web/error.tpl
dnews/web/follow.tpl
dnews/web/frame.htm
dnews/web/grpsrch.tpl
dnews/web/item.htm
dnews/web/item.tpl
dnews/web/item_f.tpl
dnews/web/item_nomore.tpl
dnews/web/items.htm
dnews/web/items.tpl
dnews/web/items3.tpl
dnews/web/items_f.tpl
dnews/web/items_t.tpl
dnews/web/itemslst.tpl
dnews/web/latest.tpl
dnews/web/list.tpl
dnews/web/login.tpl
dnews/web/missing.tpl
dnews/web/post.htm
dnews/web/post.tpl
dnews/web/post2.tpl
dnews/web/posted.tpl
dnews/web/postform.tpl
dnews/web/postit.tpl
dnews/web/required.tpl
dnews/web/results.tpl
dnews/web/results2.tpl
dnews/web/search.tpl
dnews/web/tomail_bad.tpl
dnews/web/tomail_del.tpl
dnews/web/tomail_ok.tpl
dnews/web/top.tpl
dnews/web/user.tpl
dnews/web/xresult2.tpl
dnews/web/xresults.tpl
dnews/web/xsearch.tpl
dnews/xmit
etc/dnews.conf
etc/rc.d/dnews.sh
@dirrm dnews/dmgr
@dirrm dnews/web
