# news/dnews/pkg-plist
bin/tellnews
dnews/access.conf
dnews/control.conf
dnews/dnews
dnews/dnews_start
dnews/dposter
dnews/expire.conf
dnews/fixowner.sh
dnews/fts.conf
dnews/install.txt
dnews/license.txt
dnews/manual.txt
dnews/moderators.conf
dnews/newsfeeds.conf
dnews/nocached.body
dnews/nocached.head
dnews/xmit
etc/dnews.conf
etc/rc.d/dnews.sh
