# www/ruby-cruisecontrolrb/files/cruisecontrolrb.in
#!/bin/sh
#
# PROVIDE: cruisecontrolrb
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# This script is modified by placing the following variables inside
# /etc/rc.conf.local, /etc/rc.conf, or /etc/rc.conf.d/cruisecontrolrb:
#
# cruisecontrolrb_enable (bool):            Set it to YES to enable this service.
#                                           Default: NO
# cruisecontrolrb_port (port):              Set it to a port
#                                           Default: 3333

. /etc/rc.subr

name=cruisecontrolrb
rcvar=${name}_enable

pidfile="%%PREFIX%%/www/${name}/tmp/pids/mongrel.pid"
command="%%PREFIX%%/www/${name}/cruise"

start_cmd="cruisecontrolrb_start"
stop_cmd="cruisecontrolrb_stop"

cruisecontrolrb_start() {
  set -x
    ${command} start ${command_args}
}
cruisecontrolrb_stop() {
    kill $(cat $pidfile)
}

load_rc_config $name

: ${cruisecontrolrb_enable="NO"}
: ${cruisecontrolrb_port="3333"}

command_args="-p ${cruisecontrolrb_port} -d"

run_rc_command "$1"
