# security/drweb-sendmail/pkg-plist
etc/rc.d/drwebsmf.sh
sbin/drweb-smf
share/doc/drweb-sendmail/readme.sendmail
share/doc/drweb-sendmail/readme.sendmail.rus
@dirrm share/doc/drweb-sendmail
