# www/resin3/pkg-plist
resin/doc/beta.xtp
resin/doc/changes-1.xtp
resin/doc/changes.xtp
resin/doc/errata.xtp
resin/doc/features.xtp
resin/doc/future.xtp
resin/doc/href-map.xml
resin/doc/index.xtp
resin/doc/toc.xml
resin/doc/images/caucho-blue.gif
resin/doc/images/apache_srun.gif
resin/doc/images/appserver.gif
resin/doc/images/background.gif
resin/doc/images/backup.gif
resin/doc/images/ball8.gif
resin/doc/images/bean.gif
resin/doc/images/bean1.gif
resin/doc/images/caucho.gif
resin/doc/images/compiler.gif
resin/doc/images/file.gif
resin/doc/images/folder.gif
resin/doc/images/hbleed.gif
resin/doc/images/httpd.gif
resin/doc/images/iis_srun.gif
resin/doc/images/jndi_tree.gif
resin/doc/images/jsp.gif
resin/doc/images/logo.gif
resin/doc/images/pixel.gif
resin/doc/images/resin_powered.gif
resin/doc/images/same_srun.gif
resin/doc/images/script.gif
resin/doc/images/srunc.gif
resin/doc/images/sticky-sessions.gif
resin/doc/images/style.gif
resin/doc/images/webserver.gif
resin/doc/images/website_srun.gif
resin/doc/images/caucho.jpg
resin/doc/images/.xvpics/background.gif
resin/doc/images/.xvpics/left_background.gif
resin/doc/images/caucho-white.jpg
resin/doc/images/left_background.gif
resin/lib/resin.jar
resin/lib/jsdk24.jar
resin/lib/jta_101.jar
resin/lib/isorelax.jar
resin/lib/jmx.jar
resin/lib/jca-15.jar
resin/lib/jstl-11.jar
resin/lib/webutil.jar
resin/webapps/resin-doc.war
resin/xsl/caucho/comment.xsl
resin/xsl/caucho/contents.xsl
resin/xsl/caucho/control.xsl
resin/xsl/caucho/debug.xsl
resin/xsl/caucho/defaultcopy.xsl
resin/xsl/caucho/format.xsl
resin/xsl/caucho/newstyle.xsl
resin/xsl/caucho/weblog.xsl
resin/xsl/caucho/fun.xsl
resin/xsl/caucho/href.xsl
resin/xsl/caucho/include.xsl
resin/xsl/caucho/nodeinfo.xsl
resin/xsl/caucho/toc.xsl
resin/xsl/caucho/viewfile.xsl
resin/xsl/caucho_default.xsl
resin/xsl/caucho_header.xsl
resin/xsl/caucho_subpage.xsl
resin/xsl/default.xsl
resin/xsl/doc_header.xsl
resin/xsl/header.xsl
resin/xsl/js_header.xsl
resin/xsl/new_doc_header.xsl
resin/xsl/search.xsl
resin/xsl/subpage.xsl
resin/xsl/js/caucho/defaultcopy.xsl
resin/xsl/js/caucho/format.xsl
resin/xsl/js/caucho/weblog.xsl
resin/libexec/libresin.so
resin/libexec/libresinssl.so
etc/app-default.xml
etc/rc.d/resin.sh
sbin/resinctl
%%MOD_DIR%%/mod_caucho.so
@dirrm resin/doc/images/.xvpics
@dirrm resin/doc/images
@dirrm resin/doc/css
@dirrm resin/doc/WEB-INF/classes
@dirrm resin/doc/WEB-INF/lib
@dirrm resin/doc/WEB-INF
@dirrm resin/doc
@dirrm resin/lib
@dirrm resin/webapps
@dirrm resin/xsl/caucho
@dirrm resin/xsl/js/caucho
@dirrm resin/xsl/js
@dirrm resin/xsl
@dirrm resin/libexec
@dirrm resin
@unexec if cmp -s %D/etc/resin.xml %D/etc/resin.xml-dist; then rm -f %D/etc/resin.xml; else echo If permanently deleting this package, %%PREFIX%%/etc/resin.xml must be removed manually; fi
etc/resin.xml-dist
