# multimedia/ffmpeg0/files/openrc-ffserver0.in
#!/sbin/openrc-run

name=ffserver0

command="%%PREFIX%%/bin/${name} &"
required_files=%%PREFIX%%/etc/ffserver0.conf

depend()
{
	keyword -shutdown
}
