# comms/thebridge/files/openrc-tbd.in
#!/sbin/openrc-run

name="tbd"
command="%%PREFIX%%/libexec/tbd"

depend() {
	keyword -shutdown
}
