# net/xbms/pkg-plist
@comment $FreeBSD$
bin/xbms
etc/xbms.conf.sample
etc/rc.d/xbms.sh.sample
