# net/binkd/files/binkd.in
#!/bin/sh

# PROVIDE: binkd 
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Define these binkd_* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#       /etc/rc.conf.d/apache
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#
binkd_enable=${binkd_enable:-"NO"}
binkd_config=${binkd_config:-"%%PREFIX%%/etc/binkd.conf"}
binkd_flags=${binkd_flags:-"-D"}
binkd_user=${binkd_user:-"fido"}
binkd_pidfile=${binkd_pidfile:-"/var/run/binkd/binkd.pid"}

. /etc/rc.subr
load_rc_config binkd 

name="binkd"
rcvar=binkd_enable
command="%%PREFIX%%/sbin/binkd"
command_args="${binkd_config}"
required_files="${binkd_config}"
pidfile="${binkd_pidfile}"

run_rc_command "$1"
