# net/openldap23-server/pkg-plist
bin/ldapadd
bin/ldapcompare
bin/ldapdelete
bin/ldapmodify
bin/ldapmodrdn
bin/ldappasswd
bin/ldapsearch
bin/ldapwhoami
%%SERVER%%@unexec %D/etc/rc.d/slapd.sh stop 2>/dev/null || true
%%SERVER%%@unexec %D/etc/rc.d/slurpd.sh stop 2>/dev/null || true
@unexec if cmp -s %D/etc/openldap/ldap.conf %D/etc/openldap/ldap.conf.default; then rm -f %D/etc/openldap/ldap.conf; fi
etc/openldap/ldap.conf.default
@exec [ -f %B/ldap.conf ] || cp %B/%f %B/ldap.conf
%%SERVER%%etc/openldap/schema/README
%%SERVER%%@unexec if cmp -s %D/etc/openldap/schema/corba.schema %D/etc/openldap/schema/corba.schema.default; then rm -f %D/etc/openldap/schema/corba.schema; fi
%%SERVER%%etc/openldap/schema/corba.schema.default
%%SERVER%%@exec [ -f %B/corba.schema ] || cp %B/%f %B/corba.schema
%%SERVER%%@unexec if cmp -s %D/etc/openldap/schema/core.schema %D/etc/openldap/schema/core.schema.default; then rm -f %D/etc/openldap/schema/core.schema; fi
%%SERVER%%etc/openldap/schema/core.schema.default
%%SERVER%%@exec [ -f %B/core.schema ] || cp %B/%f %B/core.schema
%%SERVER%%@unexec if cmp -s %D/etc/openldap/schema/cosine.schema %D/etc/openldap/schema/cosine.schema.default; then rm -f %D/etc/openldap/schema/cosine.schema; fi
%%SERVER%%etc/openldap/schema/cosine.schema.default
%%SERVER%%@exec [ -f %B/cosine.schema ] || cp %B/%f %B/cosine.schema
%%SERVER%%@unexec if cmp -s %D/etc/openldap/schema/inetorgperson.schema %D/etc/openldap/schema/inetorgperson.schema.default; then rm -f %D/etc/openldap/schema/inetorgperson.schema; fi
%%SERVER%%etc/openldap/schema/inetorgperson.schema.default
%%SERVER%%@exec [ -f %B/inetorgperson.schema ] || cp %B/%f %B/inetorgperson.schema
%%SERVER%%@unexec if cmp -s %D/etc/openldap/schema/java.schema %D/etc/openldap/schema/java.schema.default; then rm -f %D/etc/openldap/schema/java.schema; fi
%%SERVER%%etc/openldap/schema/java.schema.default
%%SERVER%%@exec [ -f %B/java.schema ] || cp %B/%f %B/java.schema
%%SERVER%%@unexec if cmp -s %D/etc/openldap/schema/misc.schema %D/etc/openldap/schema/misc.schema.default; then rm -f %D/etc/openldap/schema/misc.schema; fi
%%SERVER%%etc/openldap/schema/misc.schema.default
%%SERVER%%@exec [ -f %B/misc.schema ] || cp %B/%f %B/misc.schema
%%SERVER%%@unexec if cmp -s %D/etc/openldap/schema/nis.schema %D/etc/openldap/schema/nis.schema.default; then rm -f %D/etc/openldap/schema/nis.schema; fi
%%SERVER%%etc/openldap/schema/nis.schema.default
%%SERVER%%@exec [ -f %B/nis.schema ] || cp %B/%f %B/nis.schema
%%SERVER%%@unexec if cmp -s %D/etc/openldap/schema/openldap.schema %D/etc/openldap/schema/openldap.schema.default; then rm -f %D/etc/openldap/schema/openldap.schema; fi
%%SERVER%%etc/openldap/schema/openldap.schema.default
%%SERVER%%@exec [ -f %B/openldap.schema ] || cp %B/%f %B/openldap.schema
%%SERVER%%@unexec if cmp -s %D/etc/openldap/slapd.conf %D/etc/openldap/slapd.conf.default; then rm -f %D/etc/openldap/slapd.conf; fi
%%SERVER%%etc/openldap/slapd.conf.default
%%SERVER%%@exec [ -f %B/slapd.conf ] || cp %B/%f %B/slapd.conf
%%SERVER%%@unexec rmdir %D/etc/openldap/schema 2>/dev/null || true
@unexec rmdir %D/etc/openldap 2>/dev/null || true
%%SERVER%%etc/rc.d/slapd.sh
%%SERVER%%etc/rc.d/slurpd.sh
include/lber.h
include/lber_types.h
include/ldap.h
include/ldap_cdefs.h
include/ldap_features.h
include/ldap_schema.h
include/ldap_utf8.h
include/slapi-plugin.h
lib/liblber.a
lib/liblber.so
lib/liblber.so.2
lib/libldap.a
lib/libldap.so
lib/libldap.so.2
lib/libldap_r.a
lib/libldap_r.so
lib/libldap_r.so.2
%%SLAPI%%lib/libslapi.a
%%SLAPI%%lib/libslapi.so
%%SLAPI%%lib/libslapi.so.2
%%SERVER%%libexec/slapd
%%SERVER%%libexec/slurpd
%%SERVER%%sbin/slapadd
%%SERVER%%sbin/slapcat
%%SERVER%%sbin/slapindex
%%SERVER%%sbin/slappasswd
share/openldap/ucdata/case.dat
share/openldap/ucdata/cmbcl.dat
share/openldap/ucdata/comp.dat
share/openldap/ucdata/ctype.dat
share/openldap/ucdata/decomp.dat
share/openldap/ucdata/kdecomp.dat
share/openldap/ucdata/num.dat
@dirrm share/openldap/ucdata
@dirrm share/openldap
%%SERVER%%@unexec rmdir %%LOCALSTATEDIR%%/openldap-slurp 2>/dev/null || true
%%SERVER%%@unexec rmdir %%LOCALSTATEDIR%%/openldap-data 2>/dev/null || true
%%SERVER%%@unexec rmdir %%LDAP_RUN_DIR%% 2>/dev/null || true
