# mail/spamass-milter/pkg-plist
sbin/spamass-milter
etc/rc.d/spamass-milter.sh
%%PORTDOCS%%%%DOCSDIR%%/activation.txt
