# net/cvsup-mirror/pkg-plist
etc/cvsup/config.sh
etc/cvsup/cvsupd.access
etc/cvsup/sup.client/distrib/refuse.self
etc/cvsup/supfile
etc/cvsup/update.sh
etc/rc.d/cvsupd.sh
