# net-mgmt/send/pkg-plist
@comment $FreeBSD$
sbin/cgatool
sbin/ipexttool
sbin/sendd
etc/rc.d/sendd
etc/sendd/sendd.conf-dist
etc/sendd/params.conf-dist
etc/sendd/ipext.conf-dist
etc/sendd/ipext_verify.conf-dist
@dirrmtry etc/sendd
