# irc/hybserv/pkg-plist
etc/rc.d/hybserv.sh
@owner hybserv
@group hybserv
hybserv/cleandb
hybserv/doc/COPYING
hybserv/doc/INSTALL
hybserv/doc/README
hybserv/doc/README-FLAGS
hybserv/doc/RELEASENOTES
hybserv/doc/TCM-LINKING
hybserv/doc/UPGRADE
hybserv/doc/cleandb.8
hybserv/doc/hybserv.8
hybserv/doc/showchans.8
hybserv/doc/shownicks.8
hybserv/encryptconf
hybserv/encryptdb
hybserv/finddupe
hybserv/fixalvl
hybserv/glines.conf
hybserv/help/chanserv/access/add
hybserv/help/chanserv/access/del
hybserv/help/chanserv/access/index
hybserv/help/chanserv/access/list
hybserv/help/chanserv/akick/add
hybserv/help/chanserv/akick/del
hybserv/help/chanserv/akick/index
hybserv/help/chanserv/akick/list
hybserv/help/chanserv/clear
hybserv/help/chanserv/clearnoexp
hybserv/help/chanserv/cycle
hybserv/help/chanserv/delme
hybserv/help/chanserv/drop
hybserv/help/chanserv/forbid
hybserv/help/chanserv/forget
hybserv/help/chanserv/halfop
hybserv/help/chanserv/identify
hybserv/help/chanserv/index
hybserv/help/chanserv/info
hybserv/help/chanserv/invite
hybserv/help/chanserv/level/index
hybserv/help/chanserv/level/list
hybserv/help/chanserv/level/reset
hybserv/help/chanserv/level/set
hybserv/help/chanserv/list
hybserv/help/chanserv/modes
hybserv/help/chanserv/noexpire
hybserv/help/chanserv/op
hybserv/help/chanserv/register
hybserv/help/chanserv/resetlevels
hybserv/help/chanserv/set/comment
hybserv/help/chanserv/set/email
hybserv/help/chanserv/set/entrymsg
hybserv/help/chanserv/set/expirebans
hybserv/help/chanserv/set/forget
hybserv/help/chanserv/set/founder
hybserv/help/chanserv/set/guard
hybserv/help/chanserv/set/index
hybserv/help/chanserv/set/mlock
hybserv/help/chanserv/set/password
hybserv/help/chanserv/set/private
hybserv/help/chanserv/set/pubcommands
hybserv/help/chanserv/set/restricted
hybserv/help/chanserv/set/secure
hybserv/help/chanserv/set/secureops
hybserv/help/chanserv/set/splitops
hybserv/help/chanserv/set/successor
hybserv/help/chanserv/set/topic
hybserv/help/chanserv/set/topiclock
hybserv/help/chanserv/set/url
hybserv/help/chanserv/set/verbose
hybserv/help/chanserv/setpass
hybserv/help/chanserv/status
hybserv/help/chanserv/unban
hybserv/help/chanserv/unforbid
hybserv/help/chanserv/voice
hybserv/help/global/gchannote
hybserv/help/global/gcmemo
hybserv/help/global/gmemo
hybserv/help/global/gnote
hybserv/help/global/index
hybserv/help/global/motd
hybserv/help/helpserv/chanserv
hybserv/help/helpserv/global
hybserv/help/helpserv/index
hybserv/help/helpserv/memoserv
hybserv/help/helpserv/nickserv
hybserv/help/helpserv/operserv
hybserv/help/helpserv/seenserv
hybserv/help/helpserv/statserv
hybserv/help/memoserv/del
hybserv/help/memoserv/forward
hybserv/help/memoserv/index
hybserv/help/memoserv/list
hybserv/help/memoserv/purge
hybserv/help/memoserv/read
hybserv/help/memoserv/reply
hybserv/help/memoserv/send
hybserv/help/memoserv/undel
hybserv/help/nickserv/access/add
hybserv/help/nickserv/access/del
hybserv/help/nickserv/access/index
hybserv/help/nickserv/access/list
hybserv/help/nickserv/clearnoexp
hybserv/help/nickserv/collide
hybserv/help/nickserv/drop
hybserv/help/nickserv/droplink
hybserv/help/nickserv/flag
hybserv/help/nickserv/forbid
hybserv/help/nickserv/ghost
hybserv/help/nickserv/identify
hybserv/help/nickserv/index
hybserv/help/nickserv/info
hybserv/help/nickserv/link
hybserv/help/nickserv/list
hybserv/help/nickserv/noexpire
hybserv/help/nickserv/nolink
hybserv/help/nickserv/recover
hybserv/help/nickserv/recoverpass
hybserv/help/nickserv/register
hybserv/help/nickserv/release
hybserv/help/nickserv/set/automask
hybserv/help/nickserv/set/email
hybserv/help/nickserv/set/gsm
hybserv/help/nickserv/set/hide
hybserv/help/nickserv/set/index
hybserv/help/nickserv/set/kill
hybserv/help/nickserv/set/master
hybserv/help/nickserv/set/memos
hybserv/help/nickserv/set/neverop
hybserv/help/nickserv/set/noadd
hybserv/help/nickserv/set/nolink
hybserv/help/nickserv/set/notify
hybserv/help/nickserv/set/password
hybserv/help/nickserv/set/phone
hybserv/help/nickserv/set/phrase
hybserv/help/nickserv/set/private
hybserv/help/nickserv/set/privmsg
hybserv/help/nickserv/set/secure
hybserv/help/nickserv/set/signon
hybserv/help/nickserv/set/uin
hybserv/help/nickserv/set/unsecure
hybserv/help/nickserv/set/url
hybserv/help/nickserv/setpass
hybserv/help/nickserv/showlink
hybserv/help/nickserv/unforbid
hybserv/help/nickserv/unlink
hybserv/help/operserv/channel
hybserv/help/operserv/clones
hybserv/help/operserv/dcc/boot
hybserv/help/operserv/dcc/index
hybserv/help/operserv/dcc/link
hybserv/help/operserv/dcc/motd
hybserv/help/operserv/dcc/quit
hybserv/help/operserv/dcc/unlink
hybserv/help/operserv/dcc/who
hybserv/help/operserv/die
hybserv/help/operserv/dump
hybserv/help/operserv/fuckover
hybserv/help/operserv/gline
hybserv/help/operserv/help
hybserv/help/operserv/htm
hybserv/help/operserv/hub
hybserv/help/operserv/identify
hybserv/help/operserv/ignore
hybserv/help/operserv/index
hybserv/help/operserv/join
hybserv/help/operserv/jump
hybserv/help/operserv/jupe
hybserv/help/operserv/kill
hybserv/help/operserv/killchan
hybserv/help/operserv/killhost
hybserv/help/operserv/kline
hybserv/help/operserv/logout
hybserv/help/operserv/motd
hybserv/help/operserv/off
hybserv/help/operserv/omode
hybserv/help/operserv/on
hybserv/help/operserv/part
hybserv/help/operserv/rehash
hybserv/help/operserv/reload
hybserv/help/operserv/restart
hybserv/help/operserv/save
hybserv/help/operserv/secure
hybserv/help/operserv/set
hybserv/help/operserv/stats
hybserv/help/operserv/status
hybserv/help/operserv/trace
hybserv/help/operserv/umode
hybserv/help/operserv/ungline
hybserv/help/operserv/unjupe
hybserv/help/seenserv/index
hybserv/help/seenserv/join
hybserv/help/seenserv/part
hybserv/help/seenserv/seen
hybserv/help/seenserv/seennick
hybserv/help/seenserv/seenstat
hybserv/help/statserv/clearstats
hybserv/help/statserv/domain
hybserv/help/statserv/greplog
hybserv/help/statserv/host
hybserv/help/statserv/index
hybserv/help/statserv/refresh
hybserv/help/statserv/server
hybserv/help/statserv/showadmins
hybserv/help/statserv/showopers
hybserv/help/statserv/showstats
hybserv/help/statserv/split
hybserv/help/statserv/stats
hybserv/hybserv
hybserv/hybserv.conf
hybserv/jupes.conf
hybserv/mkpasswd
hybserv/motd.dcc
hybserv/motd.global
hybserv/servchk
hybserv/settings.conf
hybserv/showchans
hybserv/shownicks
@dir hybserv/logs
@owner
@group
