# net/motsognir/files/openrc-motsognir.in
#!/sbin/openrc-run
#
# motsognir.sh for rc.d usage (c) 2013-2014 Mateusz Viste.

# Add these lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# motsognir_config (path):     Set to %%PREFIX%%/etc/motsognir.conf
#                              by default.

name=motsognir

: ${motsognir_config="%%PREFIX%%/etc/motsognir.conf"}

command=%%PREFIX%%/bin/${name}
command_args="--config $motsognir_config"

depend() {
	keyword -shutdown
}
