# net-p2p/bnbt/files/bnbt.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: bnbt
# REQUIRE: LOGIN abi
# BEFORE:  securelevel
# KEYWORD: shutdown

# Add the following line to /etc/rc.conf to enable `bnbt':
#
# bnbt_enable="YES"
# bnbt_flags="<set as needed>"
#
# See bnbt(1) for bnbt_flags
#

. /etc/rc.subr

name="bnbt"
rcvar=bnbt_enable

# path to your executable, might be libxec, bin, sbin, ...
command="%%PREFIX%%/bin/$name"

# extra required arguments
command_args="&"

# you can check for required_dirs and required_vars too, see rc.subr(8)
#
required_dirs="/var/log/$name"

# read settings, set default values
load_rc_config "$name"
: ${bnbt_enable="NO"}

run_rc_command "$1"
