# net/ris/files/binlsrv.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: binlsrv
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable binlsrv:
#
# binlsrv_enable="YES"
# binlsrv_infdir="/path/to/inf/files"

. /etc/rc.subr

name="binlsrv"
rcvar=binlsrv_enable

load_rc_config "${name}"

# Set defaults
: ${binlsrv_enable:="NO"}
: ${binlsrv_infdir:="%%DATADIR%%/inf"}

pidfile="%%PID_FILE%%"
command="%%PYTHON_CMD%%"
command_args="%%DATADIR%%/${name}.py -d"
required_dirs="$binlsrv_infdir"

start_precmd="binlsrv_start_precmd"
stop_postcmd="binlsrv_stop_postcmd"

binlsrv_start_precmd()
{
	%%DATADIR%%/infparser.py ${binlsrv_infdir}
}

binlsrv_stop_postcmd()
{
	rm -f ${pidfile}
}

run_rc_command "$1"
