# graphics/hpoj/pkg-plist
bin/ptal-connect
bin/ptal-device
bin/ptal-devid
bin/ptal-pml
bin/ptal-print
bin/ptal-hp
bin/hpojip-test
sbin/ptal-mlcd
sbin/ptal-printd
sbin/ptal-photod
sbin/ptal-init
sbin/ptal-cups
lib/sane/libsane-hpoj.so.1
lib/libhpojip.so
lib/libhpojip.so.0
lib/libptal.so
lib/libptal.so.0
lib/libsane-hpoj.so
lib/libsane-hpoj.so.1
include/hpojip.h
include/ptal.h
etc/rc.d/020.ptal-init.sh.sample
%%CUPS%%libexec/cups/backend/ptal
%%PORTDOCS%%share/doc/hpoj/COPYING
%%PORTDOCS%%share/doc/hpoj/LICENSE
%%PORTDOCS%%share/doc/hpoj/LICENSE.OpenSSL
%%PORTDOCS%%share/doc/hpoj/README
%%PORTDOCS%%share/doc/hpoj/index.html
%%PORTDOCS%%share/doc/hpoj/info-devname.html
%%PORTDOCS%%share/doc/hpoj/info-pml-objects.html
%%PORTDOCS%%share/doc/hpoj/info-protocols.html
%%PORTDOCS%%share/doc/hpoj/libptal.html
%%PORTDOCS%%share/doc/hpoj/libsane-hpoj.html
%%PORTDOCS%%share/doc/hpoj/ptal-connect.html
%%PORTDOCS%%share/doc/hpoj/ptal-device.html
%%PORTDOCS%%share/doc/hpoj/ptal-devid.html
%%PORTDOCS%%share/doc/hpoj/ptal-hp.html
%%PORTDOCS%%share/doc/hpoj/ptal-init.html
%%PORTDOCS%%share/doc/hpoj/ptal-mlcd.html
%%PORTDOCS%%share/doc/hpoj/ptal-photod.html
%%PORTDOCS%%share/doc/hpoj/ptal-pml.html
%%PORTDOCS%%share/doc/hpoj/ptal-printd.html
%%PORTDOCS%%share/doc/hpoj/setup-build-details.html
%%PORTDOCS%%share/doc/hpoj/setup-build.html
%%PORTDOCS%%share/doc/hpoj/setup-connect-details.html
%%PORTDOCS%%share/doc/hpoj/setup-connect.html
%%PORTDOCS%%share/doc/hpoj/setup-photo-details.html
%%PORTDOCS%%share/doc/hpoj/setup-photo.html
%%PORTDOCS%%share/doc/hpoj/setup-print-details.html
%%PORTDOCS%%share/doc/hpoj/setup-print.html
%%PORTDOCS%%share/doc/hpoj/setup-scan-details.html
%%PORTDOCS%%share/doc/hpoj/setup-scan.html
%%PORTDOCS%%share/doc/hpoj/xojpanel.html
%%PORTDOCS%%@dirrm share/doc/hpoj
%%X11%%@cwd %%X11BASE%%
%%X11%%bin/xojpanel
