# print/hplip3/files/hpssd.in
#!/bin/sh

# PROVIDE: hpssd
# BEFORE:  cupsd
# KEYWORD: shutdown

. /etc/rc.subr

name=hpssd
rcvar=${name}_enable

command="%%PREFIX%%/sbin/${name}"
command_interpreter="%%PREFIX%%/bin/python"
pidfile="/var/run/${name}.pid"
hpssd_user=nobody
command_args=">/dev/null"

start_precmd="hpssd_prestart"

load_rc_config $name

: ${hpssd_enable="NO"}

hpssd_prestart()
{
	touch ${pidfile}
	chown ${hpssd_user} ${pidfile}
}

run_rc_command $*
