# security/skip/pkg-plist
bin/skiptool
bin/skiphost
bin/skipstat
bin/skiplog
bin/skipstat_ui
bin/skipvar
bin/test_key
bin/skipd
bin/skipdb
bin/skipca
bin/skiplocal
bin/certreq
bin/X509toHashCert
bin/dh_keygen
bin/print_cert
bin/install_skip_keys
bin/skipd_restart
bin/skipif
bin/skip_conf
etc/rc.d/skip.sh
etc/skip/skip.mkdev
etc/skip/dh_params
etc/skip/skipd.conf
etc/skip/hosts
etc/skip/SunICG_CA_selfcert
etc/skip/none.ras
etc/skip/unknown.ras
etc/skip/icon_v1.ras
etc/skip/icon_v2.ras
etc/skip/skiptool_small.ras
etc/skip/updated.ras
etc/skip/nomadic.ras
share/doc/skip/SKIP_SOFTWARE_LICENSE
share/doc/skip/BN_SOFTWARE_LICENSE
share/doc/skip/README.PATENT
share/doc/skip/00README
share/doc/skip/INSTALL
share/doc/skip/README.FreeBSD
share/doc/skip/README.FreeBSD+NAT
share/doc/skip/advanced.TOPICS
share/doc/skip/usersguide.txt
share/doc/skip/usersguide.ps
share/doc/skip/usersguide.html
share/doc/skip/inet-95.ps
modules/skip.ko
@dirrm share/doc/skip
