# mail/dbmail23/files/dbmail-imapd.sh.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: dbmail-imapd
# REQUIRE: DAEMON
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable dbmail-imapd: 
#
# dbmail_imapd_enable="YES"
# dbmail_imapd_flags="<set as needed>"
#
# See dbmail-imapd(8) for flags
#

. /etc/rc.subr

name=dbmail_imapd
rcvar=dbmail_imapd_enable

command=%%PREFIX%%/sbin/dbmail-imapd
pidfile=/var/run/dbmail-imapd.pid
required_files=%%PREFIX%%/etc/dbmail.conf

# read settings, set default values
load_rc_config "$name"
: ${dbmail_imapd_enable="NO"}

run_rc_command "$1"
