# net-mgmt/zabbix/pkg-plist
etc/zabbix/zabbix_agent.conf.sample
etc/zabbix/zabbix_agentd.conf.sample
etc/zabbix/zabbix_suckerd.conf.sample
etc/zabbix/zabbix_trapper.conf.sample
etc/zabbix/zabbix_trapperd.conf.sample
etc/rc.d/zabbix.sh.sample
etc/rc.d/zabbix-agent.sh.sample
bin/zabbix_agent
bin/zabbix_agentd
bin/zabbix_sender
bin/zabbix_suckerd
bin/zabbix_trapper
bin/zabbix_trapperd
%%DATADIR%%/php/audio/disaster_on.wav
%%DATADIR%%/php/audio/warning_off.wav
%%DATADIR%%/php/audio/warning_on.wav
%%DATADIR%%/php/about.php
%%DATADIR%%/php/actions.php
%%DATADIR%%/php/alarms.php
%%DATADIR%%/php/alerts.php
%%DATADIR%%/php/chart.php
%%DATADIR%%/php/chart2.php
%%DATADIR%%/php/chart3.php
%%DATADIR%%/php/chart4.php
%%DATADIR%%/php/chart5.php
%%DATADIR%%/php/chart_diff.php
%%DATADIR%%/php/chart_sla.php
%%DATADIR%%/php/charts.php
%%DATADIR%%/php/compare.php
%%DATADIR%%/php/config.php
%%DATADIR%%/php/css.css
%%DATADIR%%/php/graph.php
%%DATADIR%%/php/graphs.php
%%DATADIR%%/php/helpdesk.php
%%DATADIR%%/php/history.php
%%DATADIR%%/php/hosts.php
%%DATADIR%%/php/index.php
%%DATADIR%%/php/items.php
%%DATADIR%%/php/latest.php
%%DATADIR%%/php/latestalarms.php
%%DATADIR%%/php/map.php
%%DATADIR%%/php/maps.php
%%DATADIR%%/php/media.php
%%DATADIR%%/php/queue.php
%%DATADIR%%/php/report1.php
%%DATADIR%%/php/report2.php
%%DATADIR%%/php/report3.php
%%DATADIR%%/php/screenconf.php
%%DATADIR%%/php/screenedit.php
%%DATADIR%%/php/screens.php
%%DATADIR%%/php/services.php
%%DATADIR%%/php/srv_status.php
%%DATADIR%%/php/sysmap.php
%%DATADIR%%/php/sysmaps.php
%%DATADIR%%/php/tr_comments.php
%%DATADIR%%/php/tr_status.php
%%DATADIR%%/php/trend.php
%%DATADIR%%/php/trends.php
%%DATADIR%%/php/triggers.php
%%DATADIR%%/php/users.php
%%DATADIR%%/php/images/sysmaps/old/Notebook.png
%%DATADIR%%/php/images/sysmaps/old/Hub.png
%%DATADIR%%/php/images/sysmaps/old/Printer.png
%%DATADIR%%/php/images/sysmaps/old/Router.png
%%DATADIR%%/php/images/sysmaps/old/Server.png
%%DATADIR%%/php/images/sysmaps/old/Workstation.png
%%DATADIR%%/php/images/sysmaps/UPS.png
%%DATADIR%%/php/images/sysmaps/Notebook.png
%%DATADIR%%/php/images/sysmaps/Hub.png
%%DATADIR%%/php/images/sysmaps/Printer.png
%%DATADIR%%/php/images/sysmaps/Router.png
%%DATADIR%%/php/images/sysmaps/Server.png
%%DATADIR%%/php/images/sysmaps/Workstation.png
%%DATADIR%%/php/images/sysmaps/README
%%DATADIR%%/php/images/sysmaps/Phone.png
%%DATADIR%%/php/images/sysmaps/Network.png
%%DATADIR%%/php/images/sysmaps/Satellite.png
%%DATADIR%%/php/include/.htaccess
%%DATADIR%%/php/include/classes.inc.php
%%DATADIR%%/php/include/config.inc.php
%%DATADIR%%/php/include/db.inc.php
%%DATADIR%%/php/include/defines.inc.php
%%DATADIR%%/create/data/data.sql
%%DATADIR%%/create/mysql/schema.sql
%%DATADIR%%/create/postgresql/schema.sql
%%DATADIR%%/create/dbsetup.sh
%%DATADIR%%/dbpatches/1.0alpha10_to_1.0alpha11/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha10_to_1.0alpha11/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha11_to_1.0alpha12/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha11_to_1.0alpha12/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha12_to_1.0beta1/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha12_to_1.0beta1/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha3_to_1.0alpha4/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha3_to_1.0alpha4/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha4_to_1.0alpha5/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha4_to_1.0alpha5/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha6_to_1.0alpha7/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha6_to_1.0alpha7/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha7_to_1.0alpha8/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha7_to_1.0alpha8/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha8_to_1.0alpha9/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha8_to_1.0alpha9/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha9_to_1.0alpha10/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0alpha9_to_1.0alpha10/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0beta2_to_1.0beta3/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0beta2_to_1.0beta3/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0beta3_to_1.0beta4/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0beta3_to_1.0beta4/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0beta4_to_1.0beta5/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0beta4_to_1.0beta5/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0beta5_to_1.0beta6/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0beta5_to_1.0beta6/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0beta6_to_1.0beta7/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0beta6_to_1.0beta7/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0beta7_to_1.0beta8/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0beta7_to_1.0beta8/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0beta8_to_1.0beta9/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0beta8_to_1.0beta9/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0beta9_to_1.0beta10/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0beta9_to_1.0beta10/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0beta10_to_1.0beta11/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0beta10_to_1.0beta11/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0beta12_to_1.0beta13/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0beta12_to_1.0beta13/postgresql/patch.sql
%%DATADIR%%/dbpatches/1.0beta13_to_1.0beta14/mysql/patch.sql
%%DATADIR%%/dbpatches/1.0beta13_to_1.0beta14/postgresql/patch.sql
@dirrm etc/zabbix
@dirrm %%DATADIR%%/php/audio
@dirrm %%DATADIR%%/php/images/sysmaps/old
@dirrm %%DATADIR%%/php/images/sysmaps
@dirrm %%DATADIR%%/php/images
@dirrm %%DATADIR%%/php/include
@dirrm %%DATADIR%%/php
@dirrm %%DATADIR%%/create/data
@dirrm %%DATADIR%%/create/mysql
@dirrm %%DATADIR%%/create/postgresql
@dirrm %%DATADIR%%/create
@dirrm %%DATADIR%%/dbpatches/1.0alpha10_to_1.0alpha11/mysql
@dirrm %%DATADIR%%/dbpatches/1.0alpha10_to_1.0alpha11/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0alpha10_to_1.0alpha11
@dirrm %%DATADIR%%/dbpatches/1.0alpha11_to_1.0alpha12/mysql
@dirrm %%DATADIR%%/dbpatches/1.0alpha11_to_1.0alpha12/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0alpha11_to_1.0alpha12
@dirrm %%DATADIR%%/dbpatches/1.0alpha12_to_1.0beta1/mysql
@dirrm %%DATADIR%%/dbpatches/1.0alpha12_to_1.0beta1/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0alpha12_to_1.0beta1
@dirrm %%DATADIR%%/dbpatches/1.0alpha3_to_1.0alpha4/mysql
@dirrm %%DATADIR%%/dbpatches/1.0alpha3_to_1.0alpha4/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0alpha3_to_1.0alpha4
@dirrm %%DATADIR%%/dbpatches/1.0alpha4_to_1.0alpha5/mysql
@dirrm %%DATADIR%%/dbpatches/1.0alpha4_to_1.0alpha5/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0alpha4_to_1.0alpha5
@dirrm %%DATADIR%%/dbpatches/1.0alpha6_to_1.0alpha7/mysql
@dirrm %%DATADIR%%/dbpatches/1.0alpha6_to_1.0alpha7/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0alpha6_to_1.0alpha7
@dirrm %%DATADIR%%/dbpatches/1.0alpha7_to_1.0alpha8/mysql
@dirrm %%DATADIR%%/dbpatches/1.0alpha7_to_1.0alpha8/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0alpha7_to_1.0alpha8
@dirrm %%DATADIR%%/dbpatches/1.0alpha8_to_1.0alpha9/mysql
@dirrm %%DATADIR%%/dbpatches/1.0alpha8_to_1.0alpha9/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0alpha8_to_1.0alpha9
@dirrm %%DATADIR%%/dbpatches/1.0alpha9_to_1.0alpha10/mysql
@dirrm %%DATADIR%%/dbpatches/1.0alpha9_to_1.0alpha10/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0alpha9_to_1.0alpha10
@dirrm %%DATADIR%%/dbpatches/1.0beta2_to_1.0beta3/mysql
@dirrm %%DATADIR%%/dbpatches/1.0beta2_to_1.0beta3/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0beta2_to_1.0beta3
@dirrm %%DATADIR%%/dbpatches/1.0beta3_to_1.0beta4/mysql
@dirrm %%DATADIR%%/dbpatches/1.0beta3_to_1.0beta4/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0beta3_to_1.0beta4
@dirrm %%DATADIR%%/dbpatches/1.0beta4_to_1.0beta5/mysql
@dirrm %%DATADIR%%/dbpatches/1.0beta4_to_1.0beta5/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0beta4_to_1.0beta5
@dirrm %%DATADIR%%/dbpatches/1.0beta5_to_1.0beta6/mysql
@dirrm %%DATADIR%%/dbpatches/1.0beta5_to_1.0beta6/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0beta5_to_1.0beta6
@dirrm %%DATADIR%%/dbpatches/1.0beta6_to_1.0beta7/mysql
@dirrm %%DATADIR%%/dbpatches/1.0beta6_to_1.0beta7/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0beta6_to_1.0beta7
@dirrm %%DATADIR%%/dbpatches/1.0beta7_to_1.0beta8/mysql
@dirrm %%DATADIR%%/dbpatches/1.0beta7_to_1.0beta8/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0beta7_to_1.0beta8
@dirrm %%DATADIR%%/dbpatches/1.0beta8_to_1.0beta9/mysql
@dirrm %%DATADIR%%/dbpatches/1.0beta8_to_1.0beta9/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0beta8_to_1.0beta9
@dirrm %%DATADIR%%/dbpatches/1.0beta9_to_1.0beta10/mysql
@dirrm %%DATADIR%%/dbpatches/1.0beta9_to_1.0beta10/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0beta9_to_1.0beta10
@dirrm %%DATADIR%%/dbpatches/1.0beta10_to_1.0beta11/mysql
@dirrm %%DATADIR%%/dbpatches/1.0beta10_to_1.0beta11/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0beta10_to_1.0beta11
@dirrm %%DATADIR%%/dbpatches/1.0beta12_to_1.0beta13/mysql
@dirrm %%DATADIR%%/dbpatches/1.0beta12_to_1.0beta13/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0beta12_to_1.0beta13
@dirrm %%DATADIR%%/dbpatches/1.0beta13_to_1.0beta14/mysql
@dirrm %%DATADIR%%/dbpatches/1.0beta13_to_1.0beta14/postgresql
@dirrm %%DATADIR%%/dbpatches/1.0beta13_to_1.0beta14
@dirrm %%DATADIR%%/dbpatches
@dirrm %%DATADIR%%
@exec ln -sf %D/etc/zabbix /etc/zabbix
@unexec rm /etc/zabbix
