# ports-mgmt/poudriere/pkg-plist
bin/poudriere
@sample etc/poudriere.conf.sample
%%ETCDIR%%/hooks/bulk.sh.sample
%%ETCDIR%%/hooks/pkgbuild.sh.sample
@sample etc/poudriered.conf.sample
etc/rc.d/poudriered
libexec/poudriere/clock_monotonic
libexec/poudriere/cpdup
libexec/poudriere/dirempty
libexec/poudriere/dirwatch
libexec/poudriere/jexecd
libexec/poudriere/locked_mkdir
libexec/poudriere/nc
libexec/poudriere/poudriered
libexec/poudriere/rename
libexec/poudriere/rexec
libexec/poudriere/sh
libexec/poudriere/timeout
libexec/poudriere/timestamp
libexec/poudriere/tsort
man/man8/poudriere.8.gz
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/httpd.conf.sample
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/nginx.conf.sample
%%DATADIR%%/MANIFESTS
%%DATADIR%%/api.sh
%%DATADIR%%/awk/dependency_loop.awk
%%DATADIR%%/awk/humanize.awk
%%DATADIR%%/awk/json.awk
%%DATADIR%%/awk/json_jail.awk
%%DATADIR%%/awk/json_top.awk
%%DATADIR%%/awk/processonelog2.awk
%%DATADIR%%/awk/siginfo_buildtime.awk
%%DATADIR%%/bulk.sh
%%DATADIR%%/clean.sh
%%DATADIR%%/common.sh
%%DATADIR%%/daemon.sh
%%DATADIR%%/distclean.sh
%%DATADIR%%/html/assets/DataTables-1.10.0/css/jquery.dataTables.min.css
%%DATADIR%%/html/assets/DataTables-1.10.0/css/jquery.dataTables_themeroller.min.css
%%DATADIR%%/html/assets/DataTables-1.10.0/images/back_disabled.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/back_enabled.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/back_enabled_hover.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/forward_disabled.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/forward_enabled.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/forward_enabled_hover.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/sort_asc.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/sort_asc_disabled.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/sort_both.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/sort_desc.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/sort_desc_disabled.png
%%DATADIR%%/html/assets/DataTables-1.10.0/js/jquery.dataTables.min.js
%%DATADIR%%/html/assets/bootstrap-3.1.1/css/bootstrap-theme.min.css
%%DATADIR%%/html/assets/bootstrap-3.1.1/css/bootstrap.min.css
%%DATADIR%%/html/assets/bootstrap-3.1.1/fonts/glyphicons-halflings-regular.eot
%%DATADIR%%/html/assets/bootstrap-3.1.1/fonts/glyphicons-halflings-regular.svg
%%DATADIR%%/html/assets/bootstrap-3.1.1/fonts/glyphicons-halflings-regular.ttf
%%DATADIR%%/html/assets/bootstrap-3.1.1/fonts/glyphicons-halflings-regular.woff
%%DATADIR%%/html/assets/bootstrap-3.1.1/js/bootstrap.min.js
%%DATADIR%%/html/assets/favicon.ico
%%DATADIR%%/html/assets/jquery-1.11.1.js
%%DATADIR%%/html/assets/jquery-1.11.1.min.js
%%DATADIR%%/html/assets/jquery-1.11.1.min.map
%%DATADIR%%/html/assets/jquery.dataTables.rowGrouping-1.2.9.js
%%DATADIR%%/html/assets/logo.svg
%%DATADIR%%/html/assets/poudriere.css
%%DATADIR%%/html/assets/poudriere.js
%%DATADIR%%/html/build.html
%%DATADIR%%/html/index.html
%%DATADIR%%/html/jail.html
%%DATADIR%%/html/robots.txt
%%DATADIR%%/include/colors.pre.sh
%%DATADIR%%/include/colors.sh
%%DATADIR%%/include/common.sh.dragonfly
%%DATADIR%%/include/common.sh.freebsd
%%DATADIR%%/include/display.sh
%%DATADIR%%/include/fs.sh
%%DATADIR%%/include/hash.sh
%%DATADIR%%/include/html.sh
%%DATADIR%%/include/parallel.sh
%%DATADIR%%/jail.sh
%%DATADIR%%/options.sh
%%DATADIR%%/pkgclean.sh
%%DATADIR%%/ports.sh
%%DATADIR%%/processonelog.sh
%%DATADIR%%/queue.sh
%%DATADIR%%/status.sh
%%DATADIR%%/testport.sh
%%ZSH%%share/zsh/site-functions/_poudriere
