# net-im/signald/files/signald.in
#! /bin/sh -

# SPDX-License-Identifier: BSD-2-Clause
#
# Copyright (c) 2021 Michael Gmelin <grembo@FreeBSD.org>
#
# Redistribution and use in source and binary forms, with or without
# modification, are permitted provided that the following conditions are met:
#
#     1. Redistributions of source code must retain the above copyright notice,
#        this list of conditions and the following disclaimer.
#     2. Redistributions in binary form must reproduce the above copyright
#        notice, this list of conditions and the following disclaimer in the
#        documentation and/or other materials provided with the distribution.
#
# THIS SOFTWARE IS PROVIDED BY THE COPYRIGHT HOLDERS AND CONTRIBUTORS "AS IS"
# AND ANY EXPRESS OR IMPLIED WARRANTIES, INCLUDING, BUT NOT LIMITED TO, THE
# IMPLIED WARRANTIES OF MERCHANTABILITY AND FITNESS FOR A PARTICULAR PURPOSE
# ARE DISCLAIMED. IN NO EVENT SHALL THE COPYRIGHT HOLDER OR CONTRIBUTORS BE
# LIABLE FOR ANY DIRECT, INDIRECT, INCIDENTAL, SPECIAL, EXEMPLARY, OR
# CONSEQUENTIAL DAMAGES (INCLUDING, BUT NOT LIMITED TO, PROCUREMENT OF
# SUBSTITUTE GOODS OR SERVICES; LOSS OF USE, DATA, OR PROFITS; OR BUSINESS
# INTERRUPTION) HOWEVER CAUSED AND ON ANY THEORY OF LIABILITY, WHETHER IN
# CONTRACT, STRICT LIABILITY, OR TORT (INCLUDING NEGLIGENCE OR OTHERWISE)
# ARISING IN ANY WAY OUT OF THE USE OF THIS SOFTWARE, EVEN IF ADVISED OF THE
# POSSIBILITY OF SUCH DAMAGE.

# PROVIDE: signald
# REQUIRE: DAEMON LOGIN NETWORKING

. /etc/rc.subr

name='signald'
rcvar='signald_enable'

# User-facing settings and their default values.

: ${signald_enable:="NO"}
: ${signald_user:="signal-cli"}
: ${signald_group:="signal-cli"}
: ${signald_data_dir:="/var/db/signald"}
: ${signald_env:="JAVA_HOME=%%JAVA_HOME%%"}

pidfile="/var/run/${name}/${name}.pid"
procname="%%JAVA_HOME%%/bin/java"

_daemon_args="-p ${pidfile} -o /var/log/${name}/${name}.log"
_signald="%%DATADIR%%/bin/${name}"
_signald_args="-d ${signald_data_dir}"
command="/usr/sbin/daemon"
command_args="${_daemon_args} ${_signald} ${_signald_args}"

start_precmd="signald_prestart"

signald_prestart() {
	if ! install -d -o "${signald_user}" -g "${signald_group}" -m 700 \
	    "${signald_data_dir}"; then
		err 1 "Failed to create data directory \"${signald_data_dir}\""
	fi
	if ! install -d -o "${signald_user}" -g "${signald_group}" -m 755 \
	    "/var/run/${name}"; then
		err 1 "Failed to create run directory \"/var/run/${name}\""
	fi
	if ! install -d -o "${signald_user}" -g "${signald_group}" -m 755 \
	    "/var/log/${name}"; then
		err 1 "Failed to create log directory \"/var/log/${name}\""
	fi
}

load_rc_config "$name"
run_rc_command "$@"
