# astro/gpsd/files/openrc-gpsd.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf to enable gpsd:
#
# gpsd_flags (str):     Set to "" by default.
#                       See gpsd(8) for flags.
#
# gpsd_devices (str):   Set to "" by default.
#                       Example: "/dev/cuaU0" for a USB serial GPS.

name=gpsd

# Set defaults
: ${gpxlogger_flags:="-i 600"}
: ${gpxlogger_logdir:="/var/log"}
: ${gpxlogger_format:="%d %B %Y - %H:%M.gpx"}

pidfile=/var/run/gpsd.pid
command=%%PREFIX%%/sbin/gpsd
command_args="$gpsd_devices"
start_stop_daemon_args="-m"

start_post() {
	cd ${gpxlogger_logdir}
	gpx=$(date +"${gpxlogger_format}")
	touch "${gpx}"
}

start_post() {
	cd ${gpxlogger_logdir}
	gpx=$(date +"${gpxlogger_format}")
	touch "${gpx}"
}

depend(){
	use localmount
	keyword -shutdown
	provide gps
	after bootmisc
}
