# mail/proxsmtp/pkg-plist
@comment $FreeBSD$
@unexec [ ! -f /var/run/proxsmtpd.pid ] || %%PREFIX%%/etc/rc.d/proxsmtpd.sh stop || true
@unexec if cmp -s %D/etc/proxsmtpd.conf-sample %D/etc/proxsmtpd.conf; then rm -f %D/etc/proxsmtpd.conf; fi
etc/proxsmtpd.conf-sample
@exec if [ ! -f %D/etc/proxsmtpd.conf ] ; then cp -p %D/%F %B/proxsmtpd.conf; fi
sbin/proxsmtpd
