# mail/proxsmtp/pkg-plist
@unexec [ ! -f /var/run/proxsmtpd.pid ] || %%PREFIX%%/etc/rc.d/proxsmtpd.sh stop || true
@sample etc/proxsmtpd.conf.sample
man/man5/proxsmtpd.conf.5.gz
man/man8/proxsmtpd.8.gz
sbin/proxsmtpd
