# net/bpft/pkg-plist
sbin/trafd
sbin/trafstatd
bin/trafdump
bin/traflog
bin/trafsave
bin/trafstart
bin/trafstat
bin/trafstop
etc/rc.d/trafd.sh
etc/traflog.format
share/doc/bpft/README
share/doc/bpft/WhatsNew
@dirrm var/bpft
@unexec rmdir %%PREFIX%%/var 2>/dev/null || true
@dirrm share/doc/bpft
