# devel/py-mwlib/files/mwqserve.in
#! /bin/sh
#
# $FreeBSD$
#

# PROVIDE: mwqserve
# REQUIRE: DAEMON
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable mwqserve:
#
#mwqserve_enable="YES"

. /etc/rc.subr

name="mwqserve"
rcvar="mwqserve_enable"

load_rc_config $name

: ${mwqserve_user:=www}
: ${mwqserve_group:=www}
: ${mwqserve_enable:=NO}

pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args="-c -f -p ${pidfile} %%PREFIX%%/bin/mw-qserve ${mwqserve_args}"
command_interpreter="%%PYTHON_CMD%%"
procname="%%PREFIX%%/bin/mw-qserve"

start_precmd=mwqserve_pre_start

mwqserve_pre_start() {
    /usr/bin/install -o ${mwqserve_user} -g ${mwqserve_group} -m 644 /dev/null ${pidfile}
}

run_rc_command "$1"
