# devel/py-mwlib/files/postman.in
#! /bin/sh
#
# $FreeBSD$
#

# PROVIDE: postman
# REQUIRE: DAEMON
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable postman:
#
#postman_enable="YES"

. /etc/rc.subr

name="postman"
rcvar="postman_enable"

load_rc_config $name

: ${postman_user:=www}
: ${postman_group:=www}
: ${postman_enable:=NO}

pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args="-c -f -p ${pidfile} %%PREFIX%%/bin/postman ${postman_args}"
command_interpreter="%%PYTHON_CMD%%"
procname="%%PREFIX%%/bin/postman"

start_precmd=postman_pre_start

postman_pre_start() {
    /usr/bin/install -o ${postman_user} -g ${postman_group} -m 644 /dev/null ${pidfile}
}

run_rc_command "$1"
