# security/bjorb/pkg-plist
etc/bjorb.conf.sample
etc/rc.d/bjorb.sh.sample
sbin/bjorb
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog.jp
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT.jp
%%PORTDOCS%%%%DOCSDIR%%/INSTALL.jp
%%PORTDOCS%%%%DOCSDIR%%/README.jp
%%PORTDOCS%%%%DOCSDIR%%/bjorb.conf.5.jp.txt
%%PORTDOCS%%%%DOCSDIR%%/features.jp
%%PORTDOCS%%@dirrm %%DOCSDIR%%
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/bjorb.conf.doc
%%PORTEXAMPLES%%@dirrm %%EXAMPLESDIR%%
@exec ( [ -f %D/certs/bjorb.pem ] && exit ; echo "Making certificate files" ; cd %D/certs ; openssl req -new -x509 -nodes -days 365 -out bjorb.pem -keyout bjorb.pem; ln -s bjorb.pem `openssl x509 -noout -hash < bjorb.pem`.0 ; chmod 644 %D/certs/bjorb.pem )
