# security/bjorb/pkg-plist
etc/bjorb.conf.sample
etc/rc.d/bjorb.sh.sample
sbin/bjorb
share/doc/bjorb/ChangeLog
share/doc/bjorb/COPYRIGHT
share/doc/bjorb/INSTALL
share/doc/bjorb/README
share/doc/bjorb/ChangeLog.jp
share/doc/bjorb/COPYRIGHT.jp
share/doc/bjorb/INSTALL.jp
share/doc/bjorb/README.jp
share/doc/bjorb/bjorb.conf.5.jp.txt
share/doc/bjorb/features.jp
@dirrm share/doc/bjorb
share/examples/bjorb/bjorb.conf.doc
@dirrm share/examples/bjorb
@exec ( [ -f %D/certs/bjorb.pem ] && exit ; echo "Making certificate files" ; cd %D/certs ; openssl req -new -x509 -nodes -days 365 -out bjorb.pem -keyout bjorb.pem; ln -s bjorb.pem `openssl x509 -noout -hash < bjorb.pem`.0 ; chmod 644 %D/certs/bjorb.pem )
