# net/smcroute/files/openrc-smcroute.in
#!/sbin/openrc-run

name="smcroute"
command="%%PREFIX%%/sbin/${name}"
command_args="-d"

depend() {
	#need netif routing
	keyword -jail
}
