# net/bld/files/bld.in
#!/bin/sh

# PROVIDE: bld
# REQUIRE: LOGIN
# BEFORE: mail
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable bld:
#
#bld_enable="YES"
#
# See bld(8) for flags
#

. /etc/rc.subr

name="bld"
rcvar=${name}_enable

load_rc_config $name

: ${bld_enable:-"NO"}
: ${bld_flags:-"-u bld -g bld"}

command="%%PREFIX%%/sbin/bld"
extra_commands="reload"
required_dirs=/var/run/bld

run_rc_command "$1"
