# www/glewlwyd/pkg-plist
bin/glewlwyd
%%ETCDIR%%/config.json.sample
%%ETCDIR%%/glewlwyd.conf
%%ETCDIR%%/glewlwyd.conf.sample
etc/rc.d/glewlwyd
lib/glewlwyd/client/libclientmoddatabase.so
lib/glewlwyd/plugin/libprotocol_oauth2.so
lib/glewlwyd/plugin/libprotocol_oidc.so
lib/glewlwyd/plugin/libprotocol_register.so
lib/glewlwyd/scheme/libschememodcertificate.so
lib/glewlwyd/scheme/libschememodemail.so
lib/glewlwyd/scheme/libschememodhttp.so
lib/glewlwyd/scheme/libschememodoauth2.so
lib/glewlwyd/scheme/libschememodotp.so
lib/glewlwyd/scheme/libschememodpassword.so
lib/glewlwyd/user/libusermoddatabase.so
lib/glewlwyd/user/libusermodhttp.so
sbin/glewlwyd
%%DATADIR%%/CHANGELOG.md
%%DATADIR%%/README.md
%%DATADIR%%/docs/API.md
%%DATADIR%%/docs/CERTIFICATE.md
%%DATADIR%%/docs/CLIENT_DATABASE.md
%%DATADIR%%/docs/CLIENT_LDAP.md
%%DATADIR%%/docs/EMAIL.md
%%DATADIR%%/docs/GETTING_STARTED.md
%%DATADIR%%/docs/HTTP.md
%%DATADIR%%/docs/INSTALL.md
%%DATADIR%%/docs/OAUTH2.md
%%DATADIR%%/docs/OAUTH2_SCHEME.md
%%DATADIR%%/docs/OIDC.md
%%DATADIR%%/docs/OTP.md
%%DATADIR%%/docs/README.md
%%DATADIR%%/docs/REGISTER.md
%%DATADIR%%/docs/SCOPE.md
%%DATADIR%%/docs/USER.md
%%DATADIR%%/docs/USER_DATABASE.md
%%DATADIR%%/docs/USER_HTTP.md
%%DATADIR%%/docs/USER_LDAP.md
%%DATADIR%%/docs/WEBAUTHN.md
%%DATADIR%%/docs/database/README.md
%%DATADIR%%/docs/database/init-core.mariadb.sql
%%DATADIR%%/docs/database/init-core.postgre.sql
%%DATADIR%%/docs/database/init-core.sqlite3.sql
%%DATADIR%%/docs/database/init.mariadb.sql
%%DATADIR%%/docs/database/init.postgre.sql
%%DATADIR%%/docs/database/init.sqlite3.sql
%%DATADIR%%/docs/database/upgrade-2.2-core.mariadb.sql
%%DATADIR%%/docs/database/upgrade-2.2-core.postgre.sql
%%DATADIR%%/docs/database/upgrade-2.2-core.sqlite3.sql
%%DATADIR%%/docs/database/upgrade-2.3-core.mariadb.sql
%%DATADIR%%/docs/database/upgrade-2.3-core.postgre.sql
%%DATADIR%%/docs/database/upgrade-2.3-core.sqlite3.sql
%%DATADIR%%/docs/database/upgrade-2.4-core.mariadb.sql
%%DATADIR%%/docs/database/upgrade-2.4-core.postgre.sql
%%DATADIR%%/docs/database/upgrade-2.4-core.sqlite3.sql
%%DATADIR%%/docs/database/upgrade-2.5-core.mariadb.sql
%%DATADIR%%/docs/database/upgrade-2.5-core.postgre.sql
%%DATADIR%%/docs/database/upgrade-2.5-core.sqlite3.sql
%%DATADIR%%/docs/database/upgrade-2.6-core.mariadb.sql
%%DATADIR%%/docs/database/upgrade-2.6-core.postgre.sql
%%DATADIR%%/docs/database/upgrade-2.6-core.sqlite3.sql
%%DATADIR%%/docs/database/upgrade-2.7-core.mariadb.sql
%%DATADIR%%/docs/database/upgrade-2.7-core.postgre.sql
%%DATADIR%%/docs/database/upgrade-2.7-core.sqlite3.sql
%%DATADIR%%/docs/fail2ban/README.md
%%DATADIR%%/docs/fail2ban/glewlwyd-log.conf
%%DATADIR%%/docs/fail2ban/glewlwyd-syslog.conf
%%DATADIR%%/docs/fail2ban/jail.local
%%DATADIR%%/docs/glewlwyd-init
%%DATADIR%%/docs/glewlwyd.service
%%DATADIR%%/webapp/admin.js
%%DATADIR%%/webapp/callback.html
%%DATADIR%%/webapp/callback.js
%%DATADIR%%/webapp/config.json.sample
%%DATADIR%%/webapp/css/admin-custom.css
%%DATADIR%%/webapp/css/bootstrap.css
%%DATADIR%%/webapp/css/bootstrap.css.map
%%DATADIR%%/webapp/css/bootstrap.min.css
%%DATADIR%%/webapp/css/bootstrap.min.css.map
%%DATADIR%%/webapp/css/fork-awesome.css
%%DATADIR%%/webapp/css/fork-awesome.min.css
%%DATADIR%%/webapp/css/fork-awesome.min.css.map
%%DATADIR%%/webapp/css/glewlwyd-custom.css
%%DATADIR%%/webapp/css/glewlwyd.css
%%DATADIR%%/webapp/css/login-custom.css
%%DATADIR%%/webapp/css/profile-custom.css
%%DATADIR%%/webapp/css/v5-compat.css
%%DATADIR%%/webapp/css/v5-compat.min.css
%%DATADIR%%/webapp/css/v5-compat.min.css.map
%%DATADIR%%/webapp/favicon.ico
%%DATADIR%%/webapp/fonts/forkawesome-webfont.eot
%%DATADIR%%/webapp/fonts/forkawesome-webfont.svg
%%DATADIR%%/webapp/fonts/forkawesome-webfont.ttf
%%DATADIR%%/webapp/fonts/forkawesome-webfont.woff
%%DATADIR%%/webapp/fonts/forkawesome-webfont.woff2
%%DATADIR%%/webapp/img/logo-admin.png
%%DATADIR%%/webapp/img/logo-login.png
%%DATADIR%%/webapp/img/logo-profile.png
%%DATADIR%%/webapp/index.html
%%DATADIR%%/webapp/js/bootstrap.js
%%DATADIR%%/webapp/js/bootstrap.js.map
%%DATADIR%%/webapp/js/bootstrap.min.js
%%DATADIR%%/webapp/js/bootstrap.min.js.map
%%DATADIR%%/webapp/js/jquery.js
%%DATADIR%%/webapp/js/jquery.min.js
%%DATADIR%%/webapp/js/popper.js
%%DATADIR%%/webapp/js/popper.min.js
%%DATADIR%%/webapp/locales/de/translations.json
%%DATADIR%%/webapp/locales/en/translations.json
%%DATADIR%%/webapp/locales/fr/translations.json
%%DATADIR%%/webapp/locales/nl/translations.json
%%DATADIR%%/webapp/login.html
%%DATADIR%%/webapp/login.js
%%DATADIR%%/webapp/profile.html
%%DATADIR%%/webapp/profile.js
%%DATADIR%%/webapp/robots.txt
share/man/man8/glewlwyd.8.gz
%%WWWDIR%%/admin.js
%%WWWDIR%%/callback.html
%%WWWDIR%%/callback.js
%%WWWDIR%%/config.json
%%WWWDIR%%/css/admin-custom.css
%%WWWDIR%%/css/bootstrap.css
%%WWWDIR%%/css/bootstrap.css.map
%%WWWDIR%%/css/bootstrap.min.css
%%WWWDIR%%/css/bootstrap.min.css.map
%%WWWDIR%%/css/fork-awesome.css
%%WWWDIR%%/css/fork-awesome.min.css
%%WWWDIR%%/css/fork-awesome.min.css.map
%%WWWDIR%%/css/glewlwyd-custom.css
%%WWWDIR%%/css/glewlwyd.css
%%WWWDIR%%/css/login-custom.css
%%WWWDIR%%/css/profile-custom.css
%%WWWDIR%%/css/v5-compat.css
%%WWWDIR%%/css/v5-compat.min.css
%%WWWDIR%%/css/v5-compat.min.css.map
%%WWWDIR%%/favicon.ico
%%WWWDIR%%/fonts/forkawesome-webfont.eot
%%WWWDIR%%/fonts/forkawesome-webfont.svg
%%WWWDIR%%/fonts/forkawesome-webfont.ttf
%%WWWDIR%%/fonts/forkawesome-webfont.woff
%%WWWDIR%%/fonts/forkawesome-webfont.woff2
%%WWWDIR%%/img/logo-admin.png
%%WWWDIR%%/img/logo-login.png
%%WWWDIR%%/img/logo-profile.png
%%WWWDIR%%/index.html
%%WWWDIR%%/js/bootstrap.js
%%WWWDIR%%/js/bootstrap.js.map
%%WWWDIR%%/js/bootstrap.min.js
%%WWWDIR%%/js/bootstrap.min.js.map
%%WWWDIR%%/js/jquery.js
%%WWWDIR%%/js/jquery.min.js
%%WWWDIR%%/js/popper.js
%%WWWDIR%%/js/popper.min.js
%%WWWDIR%%/locales/de/translations.json
%%WWWDIR%%/locales/en/translations.json
%%WWWDIR%%/locales/fr/translations.json
%%WWWDIR%%/locales/nl/translations.json
%%WWWDIR%%/login.html
%%WWWDIR%%/login.js
%%WWWDIR%%/profile.html
%%WWWDIR%%/profile.js
%%WWWDIR%%/robots.txt
