# www/glewlwyd/pkg-plist
sbin/glewlwyd
etc/rc.d/glewlwyd
@sample %%ETCDIR%%/glewlwyd.conf.sample
@sample %%ETCDIR%%/config.json.sample
lib/glewlwyd/client/libclientmoddatabase.so
lib/glewlwyd/plugin/libprotocol_oauth2.so
lib/glewlwyd/plugin/libprotocol_oidc.so
lib/glewlwyd/plugin/libprotocol_register.so
lib/glewlwyd/scheme/libschememodcertificate.so
lib/glewlwyd/scheme/libschememodemail.so
lib/glewlwyd/scheme/libschememodhttp.so
lib/glewlwyd/scheme/libschememodotp.so
lib/glewlwyd/scheme/libschememodpassword.so
lib/glewlwyd/user/libusermoddatabase.so
lib/glewlwyd/user/libusermodhttp.so
%%DATADIR%%/CHANGELOG.md
%%DATADIR%%/README.md
%%DOCSDIR%%/API.md
%%DOCSDIR%%/CERTIFICATE.md
%%DOCSDIR%%/CLIENT_DATABASE.md
%%DOCSDIR%%/CLIENT_LDAP.md
%%DOCSDIR%%/GETTING_STARTED.md
%%DOCSDIR%%/INSTALL.md
%%DOCSDIR%%/OAUTH2.md
%%DOCSDIR%%/OIDC.md
%%DOCSDIR%%/OTP.md
%%DOCSDIR%%/README.md
%%DOCSDIR%%/SCOPE.md
%%DOCSDIR%%/USER.md
%%DOCSDIR%%/USER_DATABASE.md
%%DOCSDIR%%/USER_HTTP.md
%%DOCSDIR%%/USER_LDAP.md
%%DOCSDIR%%/WEBAUTHN.md
%%DOCSDIR%%/database/README.md
%%DOCSDIR%%/database/init-core.mariadb.sql
%%DOCSDIR%%/database/init-core.postgre.sql
%%DOCSDIR%%/database/init-core.sqlite3.sql
%%DOCSDIR%%/database/init.mariadb.sql
%%DOCSDIR%%/database/init.postgre.sql
%%DOCSDIR%%/database/init.sqlite3.sql
%%DOCSDIR%%/glewlwyd-init
%%DOCSDIR%%/glewlwyd.service
%%WWWDIR%%/admin.js
%%WWWDIR%%/config.json
%%WWWDIR%%/css/admin-custom.css
%%WWWDIR%%/css/bootstrap.css
%%WWWDIR%%/css/bootstrap.min.css
%%WWWDIR%%/css/font-awesome.css
%%WWWDIR%%/css/font-awesome.min.css
%%WWWDIR%%/css/glewlwyd-custom.css
%%WWWDIR%%/css/glewlwyd.css
%%WWWDIR%%/css/login-custom.css
%%WWWDIR%%/css/profile-custom.css
%%WWWDIR%%/favicon.ico
%%WWWDIR%%/img/logo-admin.png
%%WWWDIR%%/img/logo-login.png
%%WWWDIR%%/img/logo-profile.png
%%WWWDIR%%/index.html
%%WWWDIR%%/js/bootstrap.js
%%WWWDIR%%/js/bootstrap.min.js
%%WWWDIR%%/js/jquery.js
%%WWWDIR%%/js/jquery.min.js
%%WWWDIR%%/js/popper.js
%%WWWDIR%%/js/popper.min.js
%%WWWDIR%%/locales/en/translations.json
%%WWWDIR%%/locales/fr/translations.json
%%WWWDIR%%/locales/nl/translations.json
%%WWWDIR%%/login.html
%%WWWDIR%%/login.js
%%WWWDIR%%/profile.html
%%WWWDIR%%/profile.js
%%WWWDIR%%/vendors~admin~login~profile.js
%%WWWDIR%%/vendors~profile.js
%%WWWDIR%%/webfonts/fa-brands-400.eot
%%WWWDIR%%/webfonts/fa-brands-400.svg
%%WWWDIR%%/webfonts/fa-brands-400.ttf
%%WWWDIR%%/webfonts/fa-brands-400.woff
%%WWWDIR%%/webfonts/fa-brands-400.woff2
%%WWWDIR%%/webfonts/fa-regular-400.eot
%%WWWDIR%%/webfonts/fa-regular-400.svg
%%WWWDIR%%/webfonts/fa-regular-400.ttf
%%WWWDIR%%/webfonts/fa-regular-400.woff
%%WWWDIR%%/webfonts/fa-regular-400.woff2
%%WWWDIR%%/webfonts/fa-solid-900.eot
%%WWWDIR%%/webfonts/fa-solid-900.svg
%%WWWDIR%%/webfonts/fa-solid-900.ttf
%%WWWDIR%%/webfonts/fa-solid-900.woff
%%WWWDIR%%/webfonts/fa-solid-900.woff2
