# security/drweb-sendmail/files/drweb-smf.in
#!/bin/sh

# PROVIDE: drweb_smf
# REQUIRE: DAEMON drwebd
# BEFORE: mail

. /etc/rc.subr

name="drweb_smf"
rcvar=`set_rcvar`

load_rc_config $name

: ${drweb_smf_enable="NO"}
: ${drweb_smf_procname="%%PREFIX%%/drweb/drweb-smf"}

command=$drweb_smf_procname

run_rc_command "$1"
