# sysutils/chyves/pkg-plist
etc/rc.d/chyves
lib/chyves/chyves-basics
lib/chyves/chyves-dataset
lib/chyves/chyves-guest-console
lib/chyves/chyves-guest-creation
lib/chyves/chyves-guest-disk
lib/chyves/chyves-guest-snapshot
lib/chyves/chyves-guest-start
lib/chyves/chyves-guest-stop
lib/chyves/chyves-informational
lib/chyves/chyves-network
lib/chyves/chyves-properties
lib/chyves/chyves-resources
lib/chyves/chyves-return
lib/chyves/chyves-updates
lib/chyves/chyves-verify
man/man8/chyves.8.gz
sbin/chyves
