# mail/milter-manager/pkg-plist
bin/milter-manager-log-analyzer
bin/milter-performance-check
bin/milter-report-statistics
bin/milter-test-client
bin/milter-test-client-libmilter
bin/milter-test-server
etc/cron.d/milter-manager-log
%%ETCDIR%%/applicable-conditions/authentication.conf
%%ETCDIR%%/applicable-conditions/remote-network.conf
%%ETCDIR%%/applicable-conditions/restrict-accounts.conf
%%ETCDIR%%/applicable-conditions/s25r.conf
%%ETCDIR%%/applicable-conditions/sendmail-compatible.conf
%%ETCDIR%%/applicable-conditions/stress.conf
%%ETCDIR%%/applicable-conditions/trust.conf
%%ETCDIR%%/defaults/debian.conf
%%ETCDIR%%/defaults/freebsd.conf
%%ETCDIR%%/defaults/pkgsrc.conf
%%ETCDIR%%/defaults/redhat.conf
%%ETCDIR%%/milter-manager.conf
%%ETCDIR%%/models/example.rb
include/milter-manager/libmilter/libmilter/mfapi.h
include/milter-manager/libmilter/libmilter/mfdef.h
include/milter-manager/milter/client.h
include/milter-manager/milter/client/milter-client-context.h
include/milter-manager/milter/client/milter-client-enum-types.h
include/milter-manager/milter/client/milter-client-objects.h
include/milter-manager/milter/client/milter-client-private.h
include/milter-manager/milter/client/milter-client-runner.h
include/milter-manager/milter/client/milter-client-single-thread-runner.h
include/milter-manager/milter/client/milter-client.h
include/milter-manager/milter/core.h
include/milter-manager/milter/core/milter-agent.h
include/milter-manager/milter/core/milter-command-decoder.h
include/milter-manager/milter/core/milter-command-encoder.h
include/milter-manager/milter/core/milter-connection.h
include/milter-manager/milter/core/milter-decoder.h
include/milter-manager/milter/core/milter-encoder.h
include/milter-manager/milter/core/milter-enum-types.h
include/milter-manager/milter/core/milter-error-emittable.h
include/milter-manager/milter/core/milter-esmtp.h
include/milter-manager/milter/core/milter-event-loop.h
include/milter-manager/milter/core/milter-finished-emittable.h
include/milter-manager/milter/core/milter-glib-event-loop.h
include/milter-manager/milter/core/milter-headers.h
include/milter-manager/milter/core/milter-libev-event-loop.h
include/milter-manager/milter/core/milter-logger.h
include/milter-manager/milter/core/milter-macros-requests.h
include/milter-manager/milter/core/milter-memory-profile.h
include/milter-manager/milter/core/milter-message-result.h
include/milter-manager/milter/core/milter-option.h
include/milter-manager/milter/core/milter-protocol-agent.h
include/milter-manager/milter/core/milter-protocol.h
include/milter-manager/milter/core/milter-reader.h
include/milter-manager/milter/core/milter-reply-decoder.h
include/milter-manager/milter/core/milter-reply-encoder.h
include/milter-manager/milter/core/milter-reply-signals.h
include/milter-manager/milter/core/milter-session-result.h
include/milter-manager/milter/core/milter-syslog-logger.h
include/milter-manager/milter/core/milter-utils.h
include/milter-manager/milter/core/milter-version.h
include/milter-manager/milter/core/milter-writer.h
include/milter-manager/milter/manager.h
include/milter-manager/milter/manager/milter-manager-applicable-condition.h
include/milter-manager/milter/manager/milter-manager-child.h
include/milter-manager/milter/manager/milter-manager-children.h
include/milter-manager/milter/manager/milter-manager-configuration.h
include/milter-manager/milter/manager/milter-manager-control-command-decoder.h
include/milter-manager/milter/manager/milter-manager-control-command-encoder.h
include/milter-manager/milter/manager/milter-manager-control-protocol.h
include/milter-manager/milter/manager/milter-manager-control-reply-decoder.h
include/milter-manager/milter/manager/milter-manager-control-reply-encoder.h
include/milter-manager/milter/manager/milter-manager-controller-context.h
include/milter-manager/milter/manager/milter-manager-controller.h
include/milter-manager/milter/manager/milter-manager-egg.h
include/milter-manager/milter/manager/milter-manager-enum-types.h
include/milter-manager/milter/manager/milter-manager-launch-command-decoder.h
include/milter-manager/milter/manager/milter-manager-launch-command-encoder.h
include/milter-manager/milter/manager/milter-manager-launch-protocol.h
include/milter-manager/milter/manager/milter-manager-leader.h
include/milter-manager/milter/manager/milter-manager-module-impl.h
include/milter-manager/milter/manager/milter-manager-module.h
include/milter-manager/milter/manager/milter-manager-objects.h
include/milter-manager/milter/manager/milter-manager-process-launcher.h
include/milter-manager/milter/manager/milter-manager-reply-decoder.h
include/milter-manager/milter/manager/milter-manager-reply-encoder.h
include/milter-manager/milter/manager/milter-manager-reply-protocol.h
include/milter-manager/milter/manager/milter-manager.h
include/milter-manager/milter/server.h
include/milter-manager/milter/server/milter-server-context.h
include/milter-manager/milter/server/milter-server-enum-types.h
lib/libmilter-client.so
lib/libmilter-client.so.0
lib/libmilter-client.so.0.0.0
lib/libmilter-core.so
lib/libmilter-core.so.0
lib/libmilter-core.so.0.0.0
lib/libmilter-manager.so
lib/libmilter-manager.so.0
lib/libmilter-manager.so.0.0.0
lib/libmilter-server.so
lib/libmilter-server.so.0
lib/libmilter-server.so.0.0.0
lib/milter-manager/binding/ext/milter_client.so
lib/milter-manager/binding/ext/milter_core.so
lib/milter-manager/binding/ext/milter_manager.so
lib/milter-manager/binding/ext/milter_server.so
lib/milter-manager/binding/ext/rb-milter-client.h
lib/milter-manager/binding/ext/rb-milter-core.h
lib/milter-manager/binding/ext/rb-milter-server.h
lib/milter-manager/binding/lib/milter.rb
lib/milter-manager/binding/lib/milter/client.rb
lib/milter-manager/binding/lib/milter/client/command-line.rb
lib/milter-manager/binding/lib/milter/client/composite-session.rb
lib/milter-manager/binding/lib/milter/client/configuration.rb
lib/milter-manager/binding/lib/milter/client/context-state.rb
lib/milter-manager/binding/lib/milter/client/envelope-address.rb
lib/milter-manager/binding/lib/milter/client/fallback-session.rb
lib/milter-manager/binding/lib/milter/client/mail-transaction-shelf.rb
lib/milter-manager/binding/lib/milter/client/session-context.rb
lib/milter-manager/binding/lib/milter/client/session.rb
lib/milter-manager/binding/lib/milter/client/testing.rb
lib/milter-manager/binding/lib/milter/core.rb
lib/milter-manager/binding/lib/milter/core/callback.rb
lib/milter-manager/binding/lib/milter/core/compatible.rb
lib/milter-manager/binding/lib/milter/core/logger.rb
lib/milter-manager/binding/lib/milter/core/macro.rb
lib/milter-manager/binding/lib/milter/core/path.rb
lib/milter-manager/binding/lib/milter/core/socket-address.rb
lib/milter-manager/binding/lib/milter/manager.rb
lib/milter-manager/binding/lib/milter/manager/address-matcher.rb
lib/milter-manager/binding/lib/milter/manager/breaker.rb
lib/milter-manager/binding/lib/milter/manager/child-context.rb
lib/milter-manager/binding/lib/milter/manager/clamav-milter-config-parser.rb
lib/milter-manager/binding/lib/milter/manager/condition-table.rb
lib/milter-manager/binding/lib/milter/manager/connection-check-context.rb
lib/milter-manager/binding/lib/milter/manager/debian-init-detector.rb
lib/milter-manager/binding/lib/milter/manager/detector.rb
lib/milter-manager/binding/lib/milter/manager/enma-socket-detector.rb
lib/milter-manager/binding/lib/milter/manager/exception.rb
lib/milter-manager/binding/lib/milter/manager/file-reader.rb
lib/milter-manager/binding/lib/milter/manager/freebsd-rc-detector.rb
lib/milter-manager/binding/lib/milter/manager/init-detector.rb
lib/milter-manager/binding/lib/milter/manager/milter-greylist-config-parser.rb
lib/milter-manager/binding/lib/milter/manager/netstat-connection-checker.rb
lib/milter-manager/binding/lib/milter/manager/opendkim-config-parser.rb
lib/milter-manager/binding/lib/milter/manager/pkgsrc-rc-detector.rb
lib/milter-manager/binding/lib/milter/manager/policy-manager.rb
lib/milter-manager/binding/lib/milter/manager/postfix-cidr-table.rb
lib/milter-manager/binding/lib/milter/manager/postfix-condition-table-parser.rb
lib/milter-manager/binding/lib/milter/manager/postfix-regexp-table.rb
lib/milter-manager/binding/lib/milter/manager/rcng-detector.rb
lib/milter-manager/binding/lib/milter/manager/redhat-detector.rb
lib/milter-manager/binding/lib/milter/manager/redhat-init-detector.rb
lib/milter-manager/binding/lib/milter/manager/redhat-upstart-detector.rb
lib/milter-manager/binding/lib/milter/server.rb
lib/milter-manager/binding/lib/milter/server/testing.rb
lib/milter-manager/libmilter.so
lib/milter-manager/libmilter.so.1
lib/milter-manager/libmilter.so.1.0.1
lib/milter-manager/module/configuration/milter-manager-ruby-configuration.so
libdata/pkgconfig/milter-client.pc
libdata/pkgconfig/milter-core.pc
libdata/pkgconfig/milter-manager.pc
libdata/pkgconfig/milter-server.pc
man/ja/man1/milter-manager-log-analyzer.1.gz
man/ja/man1/milter-manager.1.gz
man/ja/man1/milter-performance-check.1.gz
man/ja/man1/milter-report-statistics.1.gz
man/ja/man1/milter-test-client.1.gz
man/ja/man1/milter-test-server.1.gz
man/man1/milter-manager-log-analyzer.1.gz
man/man1/milter-manager.1.gz
man/man1/milter-performance-check.1.gz
man/man1/milter-report-statistics.1.gz
man/man1/milter-test-client.1.gz
man/man1/milter-test-server.1.gz
sbin/milter-manager
%%DATADIR%%/license/fdl.txt
%%DATADIR%%/license/gpl.txt
%%DATADIR%%/license/lgpl.txt
%%DATADIR%%/license/ruby.ja.txt
%%DATADIR%%/license/ruby.txt
%%DATADIR%%/munin/plugins/milter_manager_
%%DATADIR%%/munin/plugins/milter_manager_memory
%%DATADIR%%/munin/plugins/postfix_processes
%%DATADIR%%/munin/plugins/sendmail_processes
