# dns/inadyn-mt/files/openrc-inadyn-mt.in
#!/sbin/openrc-run

name="inadynmt"
command=%%PREFIX%%/bin/inadyn-mt

depend() {
	use net
	keyword -shutdown
}
