# net/radvd/files/radvd.in
#!/bin/sh

# PROVIDE: radvd
# REQUIRE: NETWORKING

#
# Add the following line to /etc/rc.conf to enable radvd:
#
# radvd_enable="YES"
#

. /etc/rc.subr

name=radvd
rcvar=radvd_enable

pidfile=/var/run/${name}.pid
command="%%PREFIX%%/sbin/${name}"

load_rc_config $name
: ${radvd_enable="NO"}
: ${radvd_config="%%PREFIX%%/etc/${name}.conf"}

required_files="${radvd_config}"

command_args="-p ${pidfile} -C ${radvd_config}"

run_rc_command "$1"
