# mail/p5-Mail-SpamAssassin-devel/pkg-plist
bin/sa-learn-nonspam
bin/spamd
bin/sa-learn-rebuild
bin/sa-forget
bin/sa-learn-spam
bin/spamc
bin/spamassassin
etc/rc.d/spamd.sh-dist
etc/mail/spamassassin/local.cf
@dirrm etc/mail/spamassassin
@dirrm etc/mail
lib/perl5/site_perl/%%PERL_VER%%/%%PERL_ARCH%%/auto/Mail/SpamAssassin/.packlist
@dirrm lib/perl5/site_perl/%%PERL_VER%%/%%PERL_ARCH%%/auto/Mail/SpamAssassin
@dirrm lib/perl5/site_perl/%%PERL_VER%%/%%PERL_ARCH%%/auto/Mail
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/Conf.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/Replier.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/Bayes.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/ConfSourceSQL.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/DBBasedAddrList.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/PerMsgStatus.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/Reporter.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/SHA1.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/AuditMessage.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/HTML.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/MailingList.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/ArchiveIterator.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/CmdLearn.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/ExposedMessage.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/PerMsgLearner.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/EncappedMIME.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/BayesStore.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/NoMailAudit.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/AutoWhitelist.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/PersistentAddrList.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/Dns.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/PhraseFreqs.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/TextCat.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/EncappedMessage.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/Message.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/EvalTests.pm
lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin/Locales.pm
@dirrm lib/perl5/site_perl/%%PERL_VER%%/Mail/SpamAssassin
@dirrm lib/perl5/site_perl/%%PERL_VER%%/Mail
%%PORTDOCS%%share/doc/p5-Mail-SpamAssassin/README.sql
%%PORTDOCS%%share/doc/p5-Mail-SpamAssassin/BUGS
%%PORTDOCS%%share/doc/p5-Mail-SpamAssassin/COPYRIGHT
%%PORTDOCS%%share/doc/p5-Mail-SpamAssassin/Changes
%%PORTDOCS%%share/doc/p5-Mail-SpamAssassin/README
%%PORTDOCS%%share/doc/p5-Mail-SpamAssassin/TRADEMARK
%%PORTDOCS%%share/doc/p5-Mail-SpamAssassin/procmailrc.example
%%PORTDOCS%%share/doc/p5-Mail-SpamAssassin/sample-nonspam.txt
%%PORTDOCS%%share/doc/p5-Mail-SpamAssassin/sample-spam.txt
%%PORTDOCS%%share/doc/p5-Mail-SpamAssassin/README.spamd
%%PORTDOCS%%share/doc/p5-Mail-SpamAssassin/README.spamd-vpopmail
%%PORTDOCS%%share/doc/p5-Mail-SpamAssassin/spamassassin.sql
@dirrm %%PORTDOCS%%share/doc/p5-Mail-SpamAssassin
share/spamassassin/10_misc.cf
share/spamassassin/20_anti_ratware.cf
share/spamassassin/20_body_tests.cf
share/spamassassin/20_compensate.cf
share/spamassassin/20_head_tests.cf
share/spamassassin/20_meta_tests.cf
share/spamassassin/20_ratware.cf
share/spamassassin/20_uri_tests.cf
share/spamassassin/23_bayes.cf
share/spamassassin/25_body_tests_es.cf
share/spamassassin/25_body_tests_pl.cf
share/spamassassin/25_head_tests_es.cf
share/spamassassin/25_head_tests_pl.cf
share/spamassassin/30_text_de.cf
share/spamassassin/30_text_es.cf
share/spamassassin/30_text_fr.cf
share/spamassassin/30_text_pl.cf
share/spamassassin/40_spam_phrases.cf
share/spamassassin/50_scores.cf
share/spamassassin/60_whitelist.cf
share/spamassassin/user_prefs.template
share/spamassassin/triplets.txt
share/spamassassin/languages
@dirrm share/spamassassin
