# net/xtradius/pkg-plist
bin/radlast
bin/radtest
bin/raduse
bin/radwho
bin/radzap
sbin/checkrad
sbin/radiusd
sbin/radwatch
etc/raddb/clients
etc/raddb/dictionary
etc/raddb/dictionary.ascend
etc/raddb/dictionary.cisco
etc/raddb/dictionary.compat
etc/raddb/dictionary.livingston
etc/raddb/dictionary.shiva
etc/raddb/dictionary.tunnel
etc/raddb/dictionary.usr
etc/raddb/dictionary.versanet
etc/raddb/dictionary.xtradius
etc/raddb/execparams
etc/raddb/hints
etc/raddb/huntgroups
etc/raddb/naslist
etc/raddb/naspasswd
etc/raddb/realms
etc/raddb/users
etc/raddb/users.pg
etc/raddb/pg_auth.pl
etc/raddb/pg_acct.pl
share/doc/radiusd-xtradius/COPYRIGHT.Cistron
share/doc/radiusd-xtradius/COPYRIGHT.Livingston
share/doc/radiusd-xtradius/ChangeLog
share/doc/radiusd-xtradius/INSTALL.OLD
share/doc/radiusd-xtradius/README
share/doc/radiusd-xtradius/README.Y2K
share/doc/radiusd-xtradius/README.cisco
share/doc/radiusd-xtradius/README.pam
share/doc/radiusd-xtradius/README.proxy
share/doc/radiusd-xtradius/README.simul
share/doc/radiusd-xtradius/README.usersfile
etc/rc.d/radiusd.sh
@exec mkdir -m 700 /var/log/radacct
@dirrm etc/raddb
@dirrm share/doc/radiusd-xtradius
@unexec rm -f /var/log/radius.log
@unexec rm -f /var/run/radiusd.pid
@unexec rmdir /var/log/radacct
@unexec echo "Leaving /var/run/radutmp, /var/log/radwtmp and /var/log/radacct alone."
