# sysutils/py-google-compute-engine/files/google_clock_skew_daemon.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: google_clock_skew_daemon
# BEFORE: LOGIN
# REQUIRE: NETWORKING syslogd google_instance_setup google_network_setup
# KEYWORD: shutdown
#
# google_clock_skew_daemon_enable (bool): Set to NO by default.
#					  Set it to YES to enable google_clock_skew_daemon.
#

. /etc/rc.subr

name="google_clock_skew_daemon"
rcvar="google_clock_skew_daemon_enable"

load_rc_config "${name}"

: ${google_clock_skew_daemon_enable:="NO"}

command_interpreter="%%PYTHON_CMD%%"
command=%%PREFIX%%/bin/${name}
pidfile=/var/run/${name}.pid
start_cmd="google_clock_skew_daemon_start"

google_clock_skew_daemon_start()
{
  check_startmsgs && echo "Starting ${name}."
  daemon -p ${pidfile} -t ${name} ${command}
}

run_rc_command "$1"
