# net/arpwatch-devel/pkg-plist
sbin/arpwatch
sbin/arpsnmp
arpwatch/arp2ethers
arpwatch/ethercodes.dat
arpwatch/d.awk
arpwatch/e.awk
arpwatch/p.awk
etc/rc.d/arpwatch.sh
