# sysutils/mcollective/files/openrc-mcollectived.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# mcollectived_flags (str):	Set as needed.  Empty by default.

name="mcollectived"
command="%%PREFIX%%/sbin/mcollectived"

depend() {
	keyword -shutdown
}
