# japanese/ebnetd/files/openrc-ebnetd.in
#!/sbin/openrc-run
#
# Add the following lines to /etc/rc.conf to enable ebnetd:
# ebnetd_flags (str):    Set to "" by default.
#                        Extra flags passed to start ebnetd.

# ebnetd
name=ebnetd
command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/ebnetd/ebnd.pid"
required_dirs="/var/run/ebnetd"
required_files=%%PREFIX%%/etc/ebnetd.conf

ebnetd_flags=${ebnetd_flags:-""}

depend() {
	keyword -shutdown
}
