# databases/percona-orchestrator/pkg-plist
bin/orchestrator-client
sbin/orchestrator
@sample %%ETCDIR%%/orchestrator.conf.sample
@comment %%DATADIR%%/bin/orchestrator-client
%%DATADIR%%/metrics/orchestrator-grafana.json
%%DATADIR%%/pseudo-gtid/bin/pseudo-gtid
%%DATADIR%%/pseudo-gtid/pseudo-gtid.sql
%%DATADIR%%/pseudo-gtid/puppet/files/etc/init.d/pseudo-gtid
%%DATADIR%%/pseudo-gtid/puppet/pseudo-gtid
%%DATADIR%%/public/bootstrap/css/bootstrap-theme.css
%%DATADIR%%/public/bootstrap/css/bootstrap-theme.css.map
%%DATADIR%%/public/bootstrap/css/bootstrap-theme.min.css
%%DATADIR%%/public/bootstrap/css/bootstrap-theme.min.css.map
%%DATADIR%%/public/bootstrap/css/bootstrap.min.css
%%DATADIR%%/public/bootstrap/css/bootstrap.min.css.map
%%DATADIR%%/public/bootstrap/fonts/glyphicons-halflings-regular.eot
%%DATADIR%%/public/bootstrap/fonts/glyphicons-halflings-regular.svg
%%DATADIR%%/public/bootstrap/fonts/glyphicons-halflings-regular.ttf
%%DATADIR%%/public/bootstrap/fonts/glyphicons-halflings-regular.woff
%%DATADIR%%/public/bootstrap/fonts/glyphicons-halflings-regular.woff2
%%DATADIR%%/public/bootstrap/js/bootstrap.min.js
%%DATADIR%%/public/bootstrap/js/npm.js
%%DATADIR%%/public/css/custom.css
%%DATADIR%%/public/css/orchestrator.css
%%DATADIR%%/public/images/ajax-loader.gif
%%DATADIR%%/public/images/booking-logo-32.png
%%DATADIR%%/public/images/booking-logo-s.png
%%DATADIR%%/public/images/keep-calm-and-let-orchestrator-handle-it-transp-m.png
%%DATADIR%%/public/images/octocat-logo-32.png
%%DATADIR%%/public/images/orchestrator-logo-32-dark.png
%%DATADIR%%/public/images/orchestrator-logo-32.png
%%DATADIR%%/public/images/outbrain-logo-32.png
%%DATADIR%%/public/images/outbrain-logo-s.png
%%DATADIR%%/public/images/tile.png
%%DATADIR%%/public/js/agent.js
%%DATADIR%%/public/js/agents.js
%%DATADIR%%/public/js/audit-failure-detection.js
%%DATADIR%%/public/js/audit-recovery-shared.js
%%DATADIR%%/public/js/audit-recovery.js
%%DATADIR%%/public/js/audit.js
%%DATADIR%%/public/js/bootbox.min.js
%%DATADIR%%/public/js/cluster-analysis-shared.js
%%DATADIR%%/public/js/cluster-pools.js
%%DATADIR%%/public/js/cluster-tree.js
%%DATADIR%%/public/js/cluster.js
%%DATADIR%%/public/js/clusters-analysis.js
%%DATADIR%%/public/js/clusters.js
%%DATADIR%%/public/js/common.js
%%DATADIR%%/public/js/corex-jquery.js
%%DATADIR%%/public/js/corex.js
%%DATADIR%%/public/js/custom.js
%%DATADIR%%/public/js/d3.v3.min.js
%%DATADIR%%/public/js/discover.js
%%DATADIR%%/public/js/instance-problems.js
%%DATADIR%%/public/js/jquery-ui.min.js
%%DATADIR%%/public/js/jquery.cookie-1.4.1.min.js
%%DATADIR%%/public/js/jquery.min.js
%%DATADIR%%/public/js/md5.js
%%DATADIR%%/public/js/orchestrator.js
%%DATADIR%%/public/js/search.js
%%DATADIR%%/public/js/seed-shared.js
%%DATADIR%%/public/js/seed.js
%%DATADIR%%/public/js/seeds.js
%%DATADIR%%/public/js/status.js
%%DATADIR%%/templates/about.tmpl
%%DATADIR%%/templates/agent.tmpl
%%DATADIR%%/templates/agent_seed_details.tmpl
%%DATADIR%%/templates/agents.tmpl
%%DATADIR%%/templates/audit.tmpl
%%DATADIR%%/templates/audit_failure_detection.tmpl
%%DATADIR%%/templates/audit_recovery.tmpl
%%DATADIR%%/templates/cluster.tmpl
%%DATADIR%%/templates/cluster_pools.tmpl
%%DATADIR%%/templates/clusters.tmpl
%%DATADIR%%/templates/clusters_analysis.tmpl
%%DATADIR%%/templates/discover.tmpl
%%DATADIR%%/templates/faq.tmpl
%%DATADIR%%/templates/home.tmpl
%%DATADIR%%/templates/keep-calm.tmpl
%%DATADIR%%/templates/layout.tmpl
%%DATADIR%%/templates/search.tmpl
%%DATADIR%%/templates/seeds.tmpl
%%DATADIR%%/templates/status.tmpl
