# net-im/pickwick/pkg-plist
bin/pickwick
bin/weller
etc/rc.d/pickwick
