# mail/pop3proxy/pkg-plist
sbin/pop3proxystat
sbin/pop3proxyd
etc/rc.d/pop3proxyd.sh
etc/pop3proxyd.conf.sample
