# news/nzbget/files/openrc-nzbget.in
#!/sbin/openrc-run

name=nzbget

command=%%PREFIX%%/bin/nzbget
command_args="-D"

depend() {
		keyword -shutdown
}
