# mail/courier/files/crontab
# $FreeBSD: /tmp/pcvs/ports/mail/courier/files/Attic/crontab,v 1.1 2002-01-15 06:17:27 dwhite Exp $
#minute	hour	mday	month	wday	command
25	*	*	*	*	%%PREFIX%%/share/courier/sqwebmail/cleancache.pl
