# sysutils/idled/pkg-plist
etc/idled.cf.template
etc/rc.d/idled%%RC_SUFX%%
libexec/idled
