# sysutils/idled/pkg-plist
etc/idled.cf.template
etc/rc.d/idled.sh
libexec/idled
