# net/corosync/files/openrc-corosync-notifyd.in
#!/sbin/openrc-run

name="corosync_notifyd"
command="%%PREFIX%%/sbin/corosync-notifyd"
command_args="-l"

depend() {
	keyword -shutdown
}
