# net/corosync/files/openrc-corosync.in
#!/sbin/openrc-run

name="corosync"
command="%%PREFIX%%/sbin/corosync"

depend() {
	provide corosync
}
