# net-mgmt/rotorouter/pkg-plist
sbin/rotorouter
etc/rc.d/rotorouter.sh
etc/rotorouter.conf-example
