# net-mgmt/rotorouter/pkg-plist
sbin/rotorouter
etc/rc.d/rotorouter.sh
@unexec if cmp -s %D/etc/rotorouter.conf %D/etc/rotorouter.conf-example; then rm -f %D/etc/rotorouter.conf; fi
etc/rotorouter.conf-example
