# www/privoxy/pkg-plist
sbin/privoxy
etc/privoxy/config
etc/privoxy/default.action
etc/privoxy/default.filter
etc/privoxy/trust
etc/privoxy/templates/blocked
etc/privoxy/templates/cgi-error-404
etc/privoxy/templates/cgi-error-bad-param
etc/privoxy/templates/cgi-error-disabled
etc/privoxy/templates/cgi-error-file
etc/privoxy/templates/cgi-error-file-read-only
etc/privoxy/templates/cgi-error-modified
etc/privoxy/templates/cgi-error-parse
etc/privoxy/templates/cgi-style.css
etc/privoxy/templates/connect-failed
etc/privoxy/templates/default
etc/privoxy/templates/edit-actions-add-url-form
etc/privoxy/templates/edit-actions-for-url
etc/privoxy/templates/edit-actions-for-url-filter
etc/privoxy/templates/edit-actions-list
etc/privoxy/templates/edit-actions-list-button
etc/privoxy/templates/edit-actions-list-section
etc/privoxy/templates/edit-actions-list-url
etc/privoxy/templates/edit-actions-remove-url-form
etc/privoxy/templates/edit-actions-url-form
etc/privoxy/templates/mod-local-help
etc/privoxy/templates/mod-support-and-service
etc/privoxy/templates/mod-title
etc/privoxy/templates/mod-unstable-warning
etc/privoxy/templates/no-such-domain
etc/privoxy/templates/show-request
etc/privoxy/templates/show-status
etc/privoxy/templates/show-status-file
etc/privoxy/templates/show-url-info
etc/privoxy/templates/show-version
etc/privoxy/templates/toggle
etc/privoxy/templates/toggle-mini
etc/privoxy/templates/untrusted
etc/rc.d/privoxy.sh
%%PORTDOCS%%share/doc/privoxy-manual/actions-file.html
%%PORTDOCS%%share/doc/privoxy-manual/configuration.html
%%PORTDOCS%%share/doc/privoxy-manual/filter-file.html
%%PORTDOCS%%share/doc/privoxy-manual/introduction.html
%%PORTDOCS%%share/doc/privoxy-manual/startup.html
%%PORTDOCS%%share/doc/privoxy-manual/appendix.html
%%PORTDOCS%%share/doc/privoxy-manual/contact.html
%%PORTDOCS%%share/doc/privoxy-manual/index.html
%%PORTDOCS%%share/doc/privoxy-manual/quickstart.html
%%PORTDOCS%%share/doc/privoxy-manual/templates.html
%%PORTDOCS%%share/doc/privoxy-manual/config.html
%%PORTDOCS%%share/doc/privoxy-manual/copyright.html
%%PORTDOCS%%share/doc/privoxy-manual/installation.html
%%PORTDOCS%%share/doc/privoxy-manual/seealso.html
%%PORTDOCS%%share/doc/privoxy-manual/upgradersnote.html
%%PORTDOCS%%@dirrm share/doc/privoxy-manual
@dirrm etc/privoxy/templates
@dirrm etc/privoxy
