# dns/rbldnsd/files/rbldnsd.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: rbldnsd
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable rbldnsd:
#
#   rbldnsd_enable="YES"
#
# See rbldnsd(8) for flags.
#

. /etc/rc.subr

name=rbldnsd
rcvar=rbldnsd_enable
load_rc_config $name

rbldnsd_enable=${rbldnsd_enable:-"NO"}
rbldnsd_flags=${rbldnsd_flags:-"-r %%PREFIX%%/etc/rbldnsd -b 127.0.0.1 bl.example.com:ip4set:example"}

command=%%PREFIX%%/sbin/rbldnsd
pidfile=/var/run/rbldnsd.pid

extra_commands=reload
start_precmd=prestart
stop_postcmd="rm -f ${pidfile}"

prestart()
{
	# enforce pidfile as first argument
	rc_flags="-p ${pidfile} ${rbldnsd_flags}"
}

run_rc_command "$1"
