# japanese/multiskkserv/pkg-plist
sbin/multiskkserv
sbin/multiskkserv-ctl
bin/skkdic-p2cdb
etc/rc.d/multiskkserv.sh
