# www/hiawatha/files/openrc-hiawatha.in
#!/sbin/openrc-run

name=hiawatha

command="%%PREFIX%%/sbin/${name}"
pidfile=${hiawatha_pidfile:-"/var/run/hiawatha.pid"}

depend() {
	keyword -shutdown
}
