# www/hiawatha/files/php-fcgi.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: php_fcgi
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following line(s) to /etc/rc.conf:
# php_fcgi_enable (bool):		Set to "NO" by default, set it
#					to "YES" to enable php_fcgi.

. /etc/rc.subr

name="php_fcgi"
rcvar=`set_rcvar`

load_rc_config $name

: ${php_fcgi_enable="NO"}

command="%%PREFIX%%/sbin/php-fcgi"
pidfile=${php_fcgi_pidfile="/var/run/php-fcgi.pid"}
procname=${php_fcgi_app="%%LOCALBASE%%/bin/php-cgi"}

stop_postcmd="${name}_stop_postcmd"

php_fcgi_stop_postcmd()
{
    rm -f ${pidfile}
}

run_rc_command "$1"
