# net-mgmt/sx/files/openrc-sxserver.in
#!/sbin/openrc-run

name="sxserver"

command="%%PREFIX%%/sbin/sxserver"
command_args="start"

stop_pre() {
    eval "${command} stop"
}

restart_pre()
{
    eval "${command} restart"
}
