# dns/knot2/files/knot.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: knot
# REQUIRE: SERVERS cleanvar
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable knot:
#
# knot_enable="YES":	Set to NO by default.
#			Set it to YES to enable knot.
# knot_config="":	Set to %%PREFIX%%/etc/knot/knot.conf
#			by default.

. /etc/rc.subr

name=knot
rcvar=knot_enable

load_rc_config ${name}

: ${knot_enable:=NO}
: ${knot_config="%%PREFIX%%/etc/knot/knot.conf"}

command=%%PREFIX%%/sbin/knotd
command_args="-c ${knot_config} -d"
control=%%PREFIX%%/sbin/knotc
pidfile=/var/run/knot/knot.pid

required_files=${knot_config}

extra_commands=reload
reload_cmd="${name}_reload"
start_precmd="${name}_prestart"

knot_prestart()
{
	if [ ! -d /var/run/knot ]; then
		install -d -o %%USERS%% -g %%GROUPS%% /var/run/knot
	fi
	if [ ! -d /var/db/knot ]; then
		install -d -o %%USERS%% -g %%GROUPS%% /var/db/knot
	fi
	${control} -c ${knot_config} checkconf
}

knot_reload()
{
	echo "Reloading ${name}."
	${control} -c ${knot_config} reload
}

run_rc_command "$1"
