# databases/dbXML/files/dbxml.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: dbxml
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# dbxml_enable (bool):   Set to NO by default.
#                        Set it to YES to enable dbxml.
#

. /etc/rc.subr

name="%%PORTNAME%%"
rcvar=${name}_enable

command=%%LOCALBASE%%/bin/dbxmlctl
pidfile=%%DBXML_HOME%%/config/${name}.pid
procname=java

load_rc_config ${name}

: ${dbxml_enable="NO"}
: ${dbxml_user="%%USER%%"}

run_rc_command "$1"
