# net-mgmt/nrpe2/files/nrpe2.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: nrpe2
# REQUIRE: NETWORKING
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable nrpe2:
# nrpe2_enable (bool):    Set to "NO" by default.
#                         Set it to "YES" to enable nrpe2.
# nrpe2_flags (str):      Set to "" by default.
# nrpe2_configfile (str): Set to "%%PREFIX%%/etc/nrpe.cfg" by default.
# nrpe2_pidfile (str):    Set to "%%NAGIOSDIR%%/nrpe2.pid" by default.
#

. /etc/rc.subr

name="nrpe2"
rcvar=nrpe2_enable

: ${nrpe2_enable=NO}
: ${nrpe2_configfile=%%PREFIX%%/etc/nrpe.cfg}
pidfile=${nrpe2_pidfile:-/var/run/nrpe2/nrpe2.pid}

command="%%PREFIX%%/sbin/nrpe2"
command_args="-d"
extra_commands="reload"
start_precmd="install -d -o ${nrpe_user:-nagios} ${pidfile%/*}"
sig_reload=HUP

load_rc_config "${name}"

required_files="${nrpe2_configfile}"
command_args="${command_args} -c ${nrpe2_configfile}"

run_rc_command "$1"
