# www/resin3/files/resin3.in
#!/bin/sh
#
# $FreeBSD: /tmp/pcvs/ports/www/resin3/files/resin3.in,v 1.1 2012-04-02 23:51:25 jgh Exp $
#

# PROVIDE: %%APP_NAME%%
# REQUIRE: LOGIN NETWORKING SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable %%APP_NAME%%:
#
# %%APP_NAME%%_enable="YES"
#


%%APP_NAME%%_enable=${%%APP_NAME%%_enable:-"NO"}
%%APP_NAME%%_user=${%%APP_NAME%%_user:"www"}
%%APP_NAME%%_group=${%%APP_NAME%%_group:"www"}

. /etc/rc.subr

name=%%APP_NAME%%
rcvar=%%APP_NAME%%_enable

load_rc_config ${name}

command="%%PREFIX%%/sbin/%%APP_NAME%%ctl"
command_args="start"
pidfile=%%PID_FILE%%

if test -n "${%%APP_NAME%%_java_version}" ; then
    JAVA_HOME=$(JAVA_VERSION="${%%APP_NAME%%_java_version}" JAVAVM_DRYRUN=1 %%LOCALBASE%%/bin/java | grep JAVA_HOME | cut -d= -f2)
    procname=$(JAVA_VERSION="${%%APP_NAME%%_java_version}" JAVAVM_DRYRUN=1 %%LOCALBASE%%/bin/java | grep JAVAVM_PROG | cut -d= -f2)
else
    JAVA_HOME=$(JAVAVM_DRYRUN=1 %%LOCALBASE%%/bin/java | grep JAVA_HOME | cut -d= -f2)
    procname=$(JAVAVM_DRYRUN=1 %%LOCALBASE%%/bin/java | grep JAVAVM_PROG | cut -d= -f2)
fi

export JAVA_HOME

run_rc_command "$1"
