# net/ptpd2/files/openrc-ptpd2.in
#!/sbin/openrc-run

name="ptpd2"
command="%%PREFIX%%/sbin/${name}"

depend() {
	use net
	keyword -jail
}
