# sysutils/boxbackup-devel/files/bbstored.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: bbstored
# REQUIRE: NETWORKING
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable bbstored:
#
# bbstored_enable="YES"

. /etc/rc.subr

name="bbstored"
rcvar=bbstored_enable

load_rc_config $name

: ${bbstored_enable:="NO"}
: ${bbstored_flags:="%%PREFIX%%/etc/box/bbstored.conf"}

pidfile=${bbstored_pidfile:-"/var/run/bbstored.pid"}

command="%%PREFIX%%/sbin/bbstored"
extra_commands="reload"

run_rc_command "$1"
