# net/samba-tng/pkg-plist
bin/smbclient
bin/testparm
bin/testprns
bin/smbrun
bin/smbstatus
bin/rpcclient
bin/make_smbcodepage
bin/make_unicodemap
bin/regedit
bin/samedit
bin/net
bin/svccontrol
bin/cmdat
bin/ntspool
bin/nmblookup
bin/make_printerdef
bin/tdbimport
bin/tdbexport
bin/tdbtool
bin/smbtar
bin/addtosmbpass
bin/convert_smbpasswd
lib/codepages/codepage.1251
lib/codepages/codepage.437
lib/codepages/codepage.737
lib/codepages/codepage.775
lib/codepages/codepage.850
lib/codepages/codepage.852
lib/codepages/codepage.857
lib/codepages/codepage.861
lib/codepages/codepage.866
lib/codepages/codepage.932
lib/codepages/codepage.936
lib/codepages/codepage.949
lib/codepages/codepage.950
lib/codepages/unicode_map.437
lib/codepages/unicode_map.737
lib/codepages/unicode_map.775
lib/codepages/unicode_map.850
lib/codepages/unicode_map.852
lib/codepages/unicode_map.857
lib/codepages/unicode_map.861
lib/codepages/unicode_map.866
lib/codepages/unicode_map.932
lib/codepages/unicode_map.936
lib/codepages/unicode_map.949
lib/codepages/unicode_map.950
lib/codepages/unicode_map.ISO8859-1
lib/codepages/unicode_map.ISO8859-13
lib/codepages/unicode_map.ISO8859-15
lib/codepages/unicode_map.ISO8859-2
lib/codepages/unicode_map.ISO8859-5
lib/codepages/unicode_map.ISO8859-7
lib/codepages/unicode_map.ISO8859-9
lib/codepages/unicode_map.KOI8-R
etc/rc.d/samba-tng.sh.sample
lib/libmsrpc.so.0
lib/libmsrpc.so
lib/libmsrpc.la
lib/libsmb.so.0
lib/libsmb.so
lib/libsmb.la
lib/libnmb.so.0
lib/libnmb.so
lib/libnmb.la
lib/libsamba.so.0
lib/libsamba.so
lib/libsamba.la
lib/libsmbpw.so.0
lib/libsmbpw.so
lib/libsmbpw.la
lib/libubiqx.so.0
lib/libubiqx.so
lib/libubiqx.la
lib/libsamrpass.so.0
lib/libsamrpass.so
lib/libsamrpass.la
lib/libsurs.so.0
lib/libsurs.so
lib/libsurs.la
sbin/smbd
sbin/lsarpcd
sbin/svcctld
sbin/spoolssd
sbin/samrd
sbin/srvsvcd
sbin/wkssvcd
sbin/browserd
sbin/netlogond
sbin/winregd
sbin/nmbd
@dirrm lib/codepages
