# sysutils/upsd/pkg-plist
etc/upsd.conf
etc/rc.d/upsd.sh
sbin/upsd
