# sysutils/jailadmin/pkg-plist
etc/jailadmin.conf.sample
etc/rc.d/zzz-jail.sh
%%SITE_PERL%%/Jail.pm
%%SNMP%%sbin/jail-snmp
sbin/jailadmin
%%SNMP%%share/snmp/mibs/JAIL-MIB.txt
