# www/ephemera/pkg-plist
bin/ephemera
bin/add.munger
bin/delete.munger
bin/extract.munger
bin/replace.munger
etc/ephemera.css
etc/ephemera.config
etc/rc.d/ephemera
man/man8/ephemera.8.gz
@unexec if cmp -s %D/etc/ephemera.config %D/etc/ephemera.config.dist; then rm -f %D/etc/ephemera.config; fi
