# www/ephemera/pkg-plist
bin/ephemera
bin/ephemera-add
bin/ephemera-delete
bin/ephemera-extract
bin/ephemera-replace
bin/ephemera-generate
etc/ephemera.css.sample
etc/rc.d/ephemera
@unexec if cmp -s %D/etc/ephemera.config %D/etc/ephemera.config.dist; then rm -f %D/etc/ephemera.config; fi
etc/ephemera.config.dist
