# www/ephemera/pkg-plist
bin/ephemera
bin/add.munger
bin/create.munger
bin/delete.munger
bin/extract.munger
bin/replace.munger
bin/edit.munger
etc/rc.d/ephemera
man/man8/ephemera.8.gz
%%DATADIR%%/ephemera.css
@sample etc/ephemera.config.sample
