# www/ephemera/pkg-plist
bin/ephemera
bin/ephemera-add
bin/ephemera-delete
bin/ephemera-extract
bin/ephemera-replace
bin/ephemera-generate
etc/ephemera.css.sample
etc/rc.d/ephemera
@sample etc/ephemera.config.sample
man/man8/ephemera.8.gz
