# mail/davmail/files/davmail.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: davmail
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf or /etc/rc.conf.local to
# enable davmail:
# davmail_enable (bool):	Set to "NO" by default.
# davmail_config (path):	Set to 
#				"%%PREFIX%%/etc/davmail.properties"
#				by default

. /etc/rc.subr

name="davmail"
rcvar=davmail_enable

load_rc_config $name

: ${davmail_enable="NO"}
: ${davmail_config="%%PREFIX%%/etc/davmail.properties"}
: ${davmail_user="davmail"}
: ${davmail_group="davmail"}

start_precmd="${name}_prestart"
start_cmd="${name}_start"

davmail_procname()
{
	JAVAVM_DRYRUN=yes "%%LOCALBASE%%/bin/java" | \
	    /usr/bin/awk -F= '{if ($1=="JAVAVM_PROG") print $2}' 2> /dev/null
}

pidfiledir="/var/run/davmail"
pidfile="${pidfiledir}/davmail.pid"
procname="$(davmail_procname)"
required_files=${davmail_config}

command="/usr/sbin/daemon"
command_args="-p ${pidfile} %%LOCALBASE%%/bin/java -cp '%%DATADIR%%/davmail.jar:%%DATADIR%%/lib/*' davmail.DavGateway '${davmail_config}'"

davmail_daemon()
{
	su -m ${davmail_user} -c "env -i PATH=${PATH} ${command} ${command_args}"
}

davmail_prestart()
{
	if [ ! -d "${pidfiledir}" ]; then
		install -d -o "${davmail_user}" -g "${davmail_group}" -m 750 "${pidfiledir}"
	fi
}

davmail_start()
{
	if [ -z "$rc_fast" -a -n "$rc_pid" ]; then
		echo 1>&2 "${name} already running? (pid=$rc_pid)."
		return 1
	fi
	echo "Starting ${name}."
	{ davmail_daemon 2>&1 | /usr/bin/logger -p daemon.info -t davmail; } &
}

run_rc_command "$1"
