# audio/sndio/files/openrc-sndiod.in
#!/sbin/openrc-run

name=sndiod
command="/usr/local/bin/sndiod"
: ${sndiod_dev="rsnd/$($SYSCTL -n hw.snd.default_unit)"}
: ${sndiod_flags="-f ${sndiod_dev} -c 0:7 -j off -s default -m mon -s monitor"}
command_args=${sndiod_flags}

depend() {
	use net
	keyword -shutdown
}
