# audio/sndio/files/openrc-sndiod.in
#!/sbin/openrc-run

name=sndiod
rcvar=sndiod_enable
command="/usr/local/bin/sndiod"

