# audio/sndio/files/sndiod.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: sndiod
# REQUIRE: NETWORKING
# BEFORE:  DAEMON
# KEYWORD: shutdown

# By default sndiod will use the audio device from
# hw.snd.default_unit.  You can override this by setting sndiod_flags.
#
# To connect to a remote sndiod use e.g.
# sndiod_flags="-f snd@remotehost/0"
#
# To use /dev/dsp5
# sndiod_flags="-f rsnd/5"

. /etc/rc.subr

name=sndiod
rcvar=sndiod_enable

load_rc_config $name

: ${sndiod_enable="NO"}
: ${sndiod_flags="-s default -m mon -s monitor"}

command="%%PREFIX%%/bin/sndiod"

run_rc_command "$1"
