# www/middleman/pkg-plist
bin/mman
@sample etc/rc.d/mman.sh.sample
@sample etc/mman.xml.sample
man/man8/mman.8.gz
%%PORTDOCS%%%%DOCSDIR%%/README.html
%%PORTDOCS%%%%DOCSDIR%%/filter.pl
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@cwd /
@exec mkdir -p %%MMAN_CACHE%%
@dirrmtry %%MMAN_CACHE%%
