# www/middleman/pkg-plist
bin/mman
etc/rc.d/mman.sh.sample
etc/mman.xml.dist
%%PORTDOCS%%%%DOCSDIR%%/README.html
%%PORTDOCS%%%%DOCSDIR%%/filter.pl
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@cwd /
@exec mkdir -p %%MMAN_CACHE%%
@dirrmtry %%MMAN_CACHE%%
