# sysutils/syslog-ng/pkg-plist
etc/rc.d/syslog-ng.sh.sample
etc/syslog-ng/syslog-ng.conf.sample
sbin/syslog-ng
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/PORTS
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/sgml/syslog-ng.dvi
%%PORTDOCS%%%%DOCSDIR%%/sgml/syslog-ng.html.tar.gz
%%PORTDOCS%%%%DOCSDIR%%/sgml/syslog-ng.ps
%%PORTDOCS%%%%DOCSDIR%%/sgml/syslog-ng.sgml
%%PORTDOCS%%%%DOCSDIR%%/sgml/syslog-ng.txt
%%PORTDOCS%%@dirrm %%DOCSDIR%%/sgml
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm etc/syslog-ng
