# textproc/elasticsearch/files/elasticsearch-plugin.in
#!/bin/sh
#
# $FreeBSD$
. /etc/rc.subr

name=elasticsearch
rcvar=elasticsearch_enable
load_rc_config $name

: ${elasticsearch_config:="%%PREFIX%%/etc/elasticsearch/elasticsearch.yml"}

ES_LIB="%%PREFIX%%/lib/elasticsearch"
ES_CLASSPATH=$ES_LIB/elasticsearch-%%PORTVERSION%%.jar:$ES_LIB/*:$ES_LIB/sigar/*

JAVA=%%LOCALBASE%%/bin/java

exec $JAVA $JAVA_OPTS -Xmx64m -Xms16m -Delasticsearch -Des.config=${elasticsearch_config} -cp $ES_CLASSPATH org.elasticsearch.plugins.PluginManager $*
