# japanese/postgresql-tcltk/pkg-plist
etc/rc.d/pgsql.sh
%%INSTALL_DIR%%/.profile
%%INSTALL_DIR%%/bin/postgres
%%INSTALL_DIR%%/bin/postmaster
%%INSTALL_DIR%%/bin/ecpg
%%INSTALL_DIR%%/bin/pg_id
%%INSTALL_DIR%%/bin/pg_version
%%INSTALL_DIR%%/bin/psql
%%INSTALL_DIR%%/bin/pg_dump
%%INSTALL_DIR%%/bin/pg_dumpall
%%INSTALL_DIR%%/bin/pg_upgrade
%%INSTALL_DIR%%/bin/pg_passwd
%%INSTALL_DIR%%/bin/cleardbdir
%%INSTALL_DIR%%/bin/createdb
%%INSTALL_DIR%%/bin/createuser
%%INSTALL_DIR%%/bin/destroydb
%%INSTALL_DIR%%/bin/destroyuser
%%INSTALL_DIR%%/bin/initdb
%%INSTALL_DIR%%/bin/initlocation
%%INSTALL_DIR%%/bin/ipcclean
%%INSTALL_DIR%%/lib/global1.bki.source
%%INSTALL_DIR%%/lib/global1.description
%%INSTALL_DIR%%/lib/local1_template1.bki.source
%%INSTALL_DIR%%/lib/local1_template1.description
%%INSTALL_DIR%%/lib/pg_hba.conf.sample
%%INSTALL_DIR%%/lib/pg_geqo.sample
%%INSTALL_DIR%%/lib/libpq.a
%%INSTALL_DIR%%/lib/libpq.so.2
%%INSTALL_DIR%%/lib/libpq.so
%%INSTALL_DIR%%/lib/libecpg.a
%%INSTALL_DIR%%/lib/libecpg.so.2.6
%%INSTALL_DIR%%/lib/libecpg.so
%%INSTALL_DIR%%/lib/libpq++.a
%%INSTALL_DIR%%/lib/libpq++.so.2
%%INSTALL_DIR%%/lib/libpq++.so
%%INSTALL_DIR%%/lib/plpgsql.so
%%INSTALL_DIR%%/include/fmgr.h
%%INSTALL_DIR%%/include/lib/dllist.h
%%INSTALL_DIR%%/include/libpq/pqcomm.h
%%INSTALL_DIR%%/include/libpq/libpq-fs.h
%%INSTALL_DIR%%/include/utils/geo_decls.h
%%INSTALL_DIR%%/include/utils/elog.h
%%INSTALL_DIR%%/include/utils/palloc.h
%%INSTALL_DIR%%/include/access/attnum.h
%%INSTALL_DIR%%/include/executor/spi.h
%%INSTALL_DIR%%/include/commands/trigger.h
%%INSTALL_DIR%%/include/os.h
%%INSTALL_DIR%%/include/config.h
%%INSTALL_DIR%%/include/c.h
%%INSTALL_DIR%%/include/postgres.h
%%INSTALL_DIR%%/include/postgres_ext.h
%%INSTALL_DIR%%/include/libpq-fe.h
%%INSTALL_DIR%%/include/libpq-int.h
%%INSTALL_DIR%%/include/ecpgerrno.h
%%INSTALL_DIR%%/include/ecpglib.h
%%INSTALL_DIR%%/include/ecpgtype.h
%%INSTALL_DIR%%/include/sqlca.h
%%INSTALL_DIR%%/include/libpq++/pgenv.h
%%INSTALL_DIR%%/include/libpq++/pgconnection.h
%%INSTALL_DIR%%/include/libpq++/pgdatabase.h
%%INSTALL_DIR%%/include/libpq++/pgtransdb.h
%%INSTALL_DIR%%/include/libpq++/pgcursordb.h
%%INSTALL_DIR%%/include/libpq++/pglobject.h
%%INSTALL_DIR%%/include/libpq++.h
%%INSTALL_DIR%%/data/base/template1/pg_proc
%%INSTALL_DIR%%/data/base/template1/pg_type
%%INSTALL_DIR%%/data/base/template1/pg_attribute
%%INSTALL_DIR%%/data/base/template1/pg_class
%%INSTALL_DIR%%/data/base/template1/pg_inherits
%%INSTALL_DIR%%/data/base/template1/pg_index
%%INSTALL_DIR%%/data/base/template1/pg_statistic
%%INSTALL_DIR%%/data/base/template1/pg_operator
%%INSTALL_DIR%%/data/base/template1/pg_opclass
%%INSTALL_DIR%%/data/base/template1/pg_am
%%INSTALL_DIR%%/data/base/template1/pg_amop
%%INSTALL_DIR%%/data/base/template1/pg_amproc
%%INSTALL_DIR%%/data/base/template1/pg_language
%%INSTALL_DIR%%/data/base/template1/pg_parg
%%INSTALL_DIR%%/data/base/template1/pg_aggregate
%%INSTALL_DIR%%/data/base/template1/pg_ipl
%%INSTALL_DIR%%/data/base/template1/pg_inheritproc
%%INSTALL_DIR%%/data/base/template1/pg_rewrite
%%INSTALL_DIR%%/data/base/template1/pg_listener
%%INSTALL_DIR%%/data/base/template1/pg_description
%%INSTALL_DIR%%/data/base/template1/pg_attribute_relid_attnam_index
%%INSTALL_DIR%%/data/base/template1/pg_attribute_relid_attnum_index
%%INSTALL_DIR%%/data/base/template1/pg_attribute_attrelid_index
%%INSTALL_DIR%%/data/base/template1/pg_proc_oid_index
%%INSTALL_DIR%%/data/base/template1/pg_proc_proname_narg_type_index
%%INSTALL_DIR%%/data/base/template1/pg_proc_prosrc_index
%%INSTALL_DIR%%/data/base/template1/pg_type_oid_index
%%INSTALL_DIR%%/data/base/template1/pg_type_typname_index
%%INSTALL_DIR%%/data/base/template1/pg_class_oid_index
%%INSTALL_DIR%%/data/base/template1/pg_class_relname_index
%%INSTALL_DIR%%/data/base/template1/pg_attrdef
%%INSTALL_DIR%%/data/base/template1/pg_attrdef_adrelid_index
%%INSTALL_DIR%%/data/base/template1/pg_relcheck
%%INSTALL_DIR%%/data/base/template1/pg_relcheck_rcrelid_index
%%INSTALL_DIR%%/data/base/template1/pg_trigger
%%INSTALL_DIR%%/data/base/template1/pg_trigger_tgrelid_index
%%INSTALL_DIR%%/data/base/template1/pg_description_objoid_index
%%INSTALL_DIR%%/data/base/template1/PG_VERSION
%%INSTALL_DIR%%/data/base/template1/pg_internal.init
%%INSTALL_DIR%%/data/base/template1/pg_user
%%INSTALL_DIR%%/data/base/template1/pg_rules
%%INSTALL_DIR%%/data/base/template1/pg_views
%%INSTALL_DIR%%/data/base/template1/pg_tables
%%INSTALL_DIR%%/data/base/template1/pg_indexes
%%INSTALL_DIR%%/data/pg_variable
%%INSTALL_DIR%%/data/pg_database
%%INSTALL_DIR%%/data/pg_shadow
%%INSTALL_DIR%%/data/pg_group
%%INSTALL_DIR%%/data/pg_log
%%INSTALL_DIR%%/data/PG_VERSION
%%INSTALL_DIR%%/data/pg_hba.conf
%%INSTALL_DIR%%/data/pg_geqo.sample
%%INSTALL_DIR%%/data/pg_pwd
%%INSTALL_DIR%%/post-install-notes
share/doc/%%INSTALL_DIR%%/src/graphics/connections.ag
share/doc/%%INSTALL_DIR%%/src/graphics/clientserver.ag
share/doc/%%INSTALL_DIR%%/src/graphics/clientserver.gif
share/doc/%%INSTALL_DIR%%/src/graphics/catalogs.gif
share/doc/%%INSTALL_DIR%%/src/graphics/connections.gif
share/doc/%%INSTALL_DIR%%/src/graphics/layout.ag
share/doc/%%INSTALL_DIR%%/src/graphics/layout.gif
share/doc/%%INSTALL_DIR%%/src/graphics/catalogs.ps
share/doc/%%INSTALL_DIR%%/src/sgml/xtypes.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/admin.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/advanced.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/arch-dev.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/arch-pg.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/arch.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/array.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/biblio.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/compiler.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/contacts.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/datatype.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/dfunc.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/docguide.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ecpg.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/environ.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/extend.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/func-ref.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/geqo.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/gist.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/inherit.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/intro.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/install.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/intro-pg.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/jdbc.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/libpgtcl.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/libpq.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/lobj.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/manage.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/odbc.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/pgaccess.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ports.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/postgres.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/programmer.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/protocol.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/psql.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/query-ug.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/query.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/recovery.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/regress.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/release.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/rules.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/spi.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/start-ag.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/start.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/storage.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/trigger.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/tutorial.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/user.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/xaggr.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/xfunc.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/xindex.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/xoper.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/Makefile
share/doc/%%INSTALL_DIR%%/src/sgml/ref/vacuum.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/allfiles.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/alter_table.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/alter_user.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/begin.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/close.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/cluster.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/commands.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/commit.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/copy.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/create_aggregate.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/create_database.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/create_function.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/create_index.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/create_language.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/create_operator.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/create_rule.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/create_sequence.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/create_table.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/create_trigger.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/create_type.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/create_user.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/create_view.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/createdb.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/createuser.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/current_date.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/current_time.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/current_timestamp.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/current_user.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/declare.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/delete.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/destroydb.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/destroyuser.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/drop_aggregate.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/drop_database.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/drop_index.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/drop_function.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/drop_language.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/drop_operator.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/drop_rule.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/drop_sequence.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/drop_table.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/drop_trigger.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/drop_type.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/drop_user.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/drop_view.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/explain.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/fetch.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/grant.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/initdb.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/initlocation.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/insert.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/listen.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/load.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/lock.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/move.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/notify.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/reset.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/pg_dump.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/pg_dumpall.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/pg_upgrade.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/psql-ref.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/revoke.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/rollback.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/select.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/set.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/show.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/unlisten.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/update.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/ref/abort.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/about.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/bki.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/config.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/current.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/oper.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/func.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/history.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/indices.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/info.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/installation.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/intro-ag.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/keys.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/legal.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/libpq++.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/notation.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/page.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/pg_options.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/reference.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/runtime.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/security.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/signals.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/syntax.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/typeconv.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/xplang.sgml
share/doc/%%INSTALL_DIR%%/src/sgml/y2k.sgml
share/doc/%%INSTALL_DIR%%/src/Makefile
share/doc/%%INSTALL_DIR%%/user.tar.gz
share/doc/%%INSTALL_DIR%%/user.ps.gz
share/doc/%%INSTALL_DIR%%/tutorial.tar.gz
share/doc/%%INSTALL_DIR%%/tutorial.ps.gz
share/doc/%%INSTALL_DIR%%/programmer.tar.gz
share/doc/%%INSTALL_DIR%%/programmer.ps.gz
share/doc/%%INSTALL_DIR%%/postgres.tar.gz
share/doc/%%INSTALL_DIR%%/bug.template
share/doc/%%INSTALL_DIR%%/admin.tar.gz
share/doc/%%INSTALL_DIR%%/admin.ps.gz
share/doc/%%INSTALL_DIR%%/TODO
share/doc/%%INSTALL_DIR%%/README.mb.jp
share/doc/%%INSTALL_DIR%%/README.mb
share/doc/%%INSTALL_DIR%%/README.fsync
share/doc/%%INSTALL_DIR%%/Makefile
share/doc/%%INSTALL_DIR%%/FAQ_Linux
share/doc/%%INSTALL_DIR%%/FAQ_Irix
share/doc/%%INSTALL_DIR%%/FAQ_FreeBSD
share/doc/%%INSTALL_DIR%%/FAQ
share/doc/%%INSTALL_DIR%%/README.locale
share/doc/%%INSTALL_DIR%%/README.inet
share/doc/%%INSTALL_DIR%%/FAQ_DEV
share/doc/%%INSTALL_DIR%%/FAQ_CVS
@exec /usr/bin/env OBJFORMAT=%%PORTOBJFORMAT%% /sbin/ldconfig -m %B
@unexec /usr/bin/env OBJFORMAT=%%PORTOBJFORMAT%% /sbin/ldconfig -R
@dirrm share/doc/%%INSTALL_DIR%%/src/graphics
@dirrm share/doc/%%INSTALL_DIR%%/src/sgml/ref
@dirrm share/doc/%%INSTALL_DIR%%/src/sgml
@dirrm share/doc/%%INSTALL_DIR%%/src
@dirrm share/doc/%%INSTALL_DIR%%
@dirrm %%INSTALL_DIR%%/bin
@dirrm %%INSTALL_DIR%%/data/base/template1
@dirrm %%INSTALL_DIR%%/data/base
@dirrm %%INSTALL_DIR%%/data
@dirrm %%INSTALL_DIR%%/include/utils
@dirrm %%INSTALL_DIR%%/include/libpq++
@dirrm %%INSTALL_DIR%%/include/libpq
@dirrm %%INSTALL_DIR%%/include/lib
@exec mkdir -p %D/%%INSTALL_DIR%%/include/port/freebsd
@dirrm %%INSTALL_DIR%%/include/port/freebsd
@dirrm %%INSTALL_DIR%%/include/port
@dirrm %%INSTALL_DIR%%/include/executor
@dirrm %%INSTALL_DIR%%/include/commands
@dirrm %%INSTALL_DIR%%/include/access
@dirrm %%INSTALL_DIR%%/include
@dirrm %%INSTALL_DIR%%/lib
@dirrm %%INSTALL_DIR%%/man/manl
@dirrm %%INSTALL_DIR%%/man/man5
@dirrm %%INSTALL_DIR%%/man/man3
@dirrm %%INSTALL_DIR%%/man/man1
@dirrm %%INSTALL_DIR%%/man
@dirrm %%INSTALL_DIR%%
