# net/trafd/pkg-plist
bin/trafd
bin/trafdump
bin/traflog
bin/trafsave
bin/trafstart
bin/trafstat
bin/trafstatd
bin/trafstop
etc/traflog.format
etc/rc.d/trafd.sh.sample
%%PORTDOCS%%share/doc/trafd/README
%%PORTDOCS%%@dirrm share/doc/trafd
@exec mkdir %D/var 2>/dev/null || true
@exec mkdir %D/var/log 2>/dev/null || true
@exec mkdir %D/var/trafd 2>/dev/null || true
@unexec rmdir %D/var/trafd 2>/dev/null || true
@unexec rmdir %D/var/log 2>/dev/null || true
@unexec rmdir %D/var 2>/dev/null || true
