# net/coda6_server/pkg-plist
bin/norton
bin/norton-reinit
bin/reinit
bin/getvolinfo
bin/rpc2ping
bin/rvmsizer
bin/smon2
etc/coda/server.conf.ex
@dirrm etc/coda
sbin/au
sbin/auth2
sbin/backup
sbin/backup.sh
sbin/bldvldb.sh
sbin/codaconfedit
sbin/coda-server-logrotate
sbin/coda-setup-ports
sbin/codasrv
sbin/codastart
sbin/createvol_rep
sbin/initpw
sbin/inoder
sbin/merge
sbin/parserecdump
sbin/partial-reinit.sh
sbin/pdbtool
sbin/printvrdb
sbin/purgevol
sbin/purgevol_rep
sbin/readdump
sbin/rpc2portmap
sbin/startserver
sbin/tape.pl
sbin/updateclnt
sbin/updatefetch
sbin/updatesrv
sbin/vice-killvolumes
sbin/vice-setup
sbin/vice-setup-rvm
sbin/vice-setup-scm
sbin/vice-setup-srvdir
sbin/vice-setup-user
sbin/volutil
etc/rc.d/rc.vice
