# net-mgmt/nagiosgrapher/files/nagiosgrapher.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: %%NAME%%
# REQUIRE: NETWORKING SERVERS
# BEFORE:  DAEMON
# KEYWORD: shutdown

# Add the following line to /etc/rc.conf to enable `%%NAME%%':
#
#%%NAME%%_enable="YES"

. /etc/rc.subr

name="%%NAME%%"
: ${%%NAME%%_user="nagios"}
rcvar=`set_rcvar`

pidfile="/var/spool/nagios/nagios_grapher/$name.pid"

command="%%PREFIX%%/libexec/nagios/contrib/$name"
command_interpreter="/usr/bin/perl"

# read configuration and set defaults
load_rc_config $name
: ${%%NAME%%_enable="NO"}

run_rc_command "$1"
