# audio/murmur/files/murmur.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: murmur
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# murmur_enable (bool):		Set to NO by default.
#				Set it to YES to enable murmur.
# murmur_flags (str):		Set to "" by default.
#				Extra flags passed to start command.
#

. /etc/rc.subr

# Set some defaults
: ${murmur_enable="NO"}
: ${murmur_user="murmur"}
: ${murmur_flags=""}

name=murmur
rcvar=${name}_enable

command="%%PREFIX%%/sbin/murmurd"
command_args="-ini %%PREFIX%%/etc/murmur.ini"
pidfile="/var/run/murmur/murmur.pid"

load_rc_config $name

run_rc_command "$1"
