# net/netams/pkg-plist
@unexec if cmp -s %D/etc/netams.conf.sample %D/etc/netams.conf; then rm -f %D/etc/netams.conf; fi
etc/netams.conf.sample
etc/rc.d/netams.sh.sample
libexec/netams
bin/netamsctl
%%PORTDOCS%%share/doc/netams/documentation-en.html
%%PORTDOCS%%share/doc/netams/diagram-1.gif
%%PORTDOCS%%share/doc/netams/logo.gif
%%PORTDOCS%%share/doc/netams/net-1-1a.gif
%%PORTDOCS%%share/doc/netams/net-1-2a.gif
%%PORTDOCS%%share/doc/netams/net-1-3a.gif
share/examples/netams/netams-nat-network.cfg
share/examples/netams/netams-netflow.cfg
share/examples/netams/netams-simple.cfg
share/examples/netams/netams.cfg
share/examples/netams/.netamsctl.rc
share/examples/netams/cgi-bin/logo.gif
share/examples/netams/cgi-bin/netams_api.pl
share/examples/netams/cgi-bin/netams_example.cgi
share/examples/netams/cgi-bin/netams_graph.cgi
share/examples/netams/cgi-bin/weblogin.cgi
share/examples/netams/cgi-bin/weblogin.tem
%%PORTDOCS%%@dirrm share/doc/netams
@dirrm share/examples/netams/cgi-bin
@dirrm share/examples/netams
