# net/traff/files/openrc-traff.in
#!/sbin/openrc-run

name=traff

command=%%PREFIX%%/sbin/traff
pidfile=/var/run/traff.pid
required_files=%%PREFIX%%/etc/traff.conf

depend() {
	keyword -shutdown
}
