# net/mosquitto/files/mosquitto.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: mosquitto
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable mosquitto:
#
# mosquitto_enable="YES"
# mosquitto_flags="<set as needed>"
#
# See mosquitto(8) for flags
#

. /etc/rc.subr

name=mosquitto
rcvar=mosquitto_enable

mosquitto_enable=${mosquitto_enable:="NO"}
mosquitto_config=${mosquitto_config:="%%PREFIX%%/etc/mosquitto/mosquitto.conf"}

command=%%PREFIX%%/sbin/mosquitto
command_args="-c ${mosquitto_config} -d"
pidfile=$(grep pidfile ${mosquitto_config} | awk '{print($2)}')
required_files=${mosquitto_config}

extra_commands="reload"
stop_postcmd=stop_postcmd
stop_postcmd()
{
  rm -f $pidfile
}

load_rc_config $name
run_rc_command "$1"
