# mail/courier-imap/pkg-plist
bin/authtest
%%OPENSSLFLAG%%bin/couriertls
bin/deliverquota
bin/imapd
bin/maildirmake
bin/pop3d
etc/courier-imap/authdaemonrc.dist
%%LDAPFLAG%%etc/courier-imap/authldaprc.dist
%%MYSQLFLAG%%etc/courier-imap/authmysqlrc.dist
%%PGSQLFLAG%%etc/courier-imap/authpgsqlrc.dist
etc/courier-imap/imapd-ssl.dist
etc/courier-imap/imapd.cnf.dist
etc/courier-imap/imapd.dist
etc/courier-imap/pop3d-ssl.dist
etc/courier-imap/pop3d.cnf.dist
etc/courier-imap/pop3d.dist
etc/courier-imap/quotawarnmsg.example
etc/rc.d/courier-imap-imapd.sh.sample
etc/rc.d/courier-imap-pop3d.sh.sample
libexec/courier-imap/authlib/authdaemon
libexec/courier-imap/authlib/authdaemond
%%LDAPFLAG%%libexec/courier-imap/authlib/authdaemond.ldap
%%MYSQLFLAG%%libexec/courier-imap/authlib/authdaemond.mysql
%%PGSQLFLAG%%libexec/courier-imap/authlib/authdaemond.pgsql
libexec/courier-imap/authlib/authdaemond.plain
libexec/courier-imap/courierlogger
libexec/courier-imap/couriertcpd
libexec/courier-imap/imapd-ssl.rc
libexec/courier-imap/imapd.rc
libexec/courier-imap/makedatprog
libexec/courier-imap/pop3d-ssl.rc
libexec/courier-imap/pop3d.rc
sbin/imaplogin
sbin/makeuserdb
sbin/mkimapdcert
sbin/mkpop3dcert
sbin/pop3login
sbin/pw2userdb
sbin/userdb
sbin/userdbpw
sbin/vchkpw2userdb
share/courier-imap/makeuserdb
share/courier-imap/mkimapdcert
share/courier-imap/mkpop3dcert
share/courier-imap/pw2userdb
share/courier-imap/userdb
share/courier-imap/vchkpw2userdb
share/doc/courier-imap/AUTHORS
share/doc/courier-imap/ChangeLog
share/doc/courier-imap/INSTALL
share/doc/courier-imap/README
@dirrm var/authdaemon
@dirrm var
@dirrm share/doc/courier-imap
@dirrm share/courier-imap
@dirrm libexec/courier-imap/authlib
@dirrm libexec/courier-imap
@dirrm etc/courier-imap
