# mail/exact/pkg-plist
@sample etc/exact.conf.sample
etc/rc.d/exact.sh
man/man5/exact.conf.5.gz
man/man8/exact.8.gz
sbin/exact
@owner %%MAILUSER%%
@group %%MAILGROUP%%
@dir %%EXACTDIR%%
@owner
@group
