# mail/exact/pkg-plist
sbin/exact
etc/exact.conf.sample
etc/rc.d/exact.sh
