# misc/py-glance/files/glance-api.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: glance_api
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable glance-api:
#
# glance_api_enable (bool):	Set it to "YES" to enable glance-api.
# 				Default is "NO".

. /etc/rc.subr

name="glance_api"
rcvar=glance_api_enable

extra_commands="reload shutdown status"
reload_cmd="glance-control api reload"
shutdown_cmd="glance-control api shutdown"
start_cmd="glance-control api start"
status_cmd="glance-control api status"
stop_cmd="glance-control api stop"

required_files=%%PREFIX%%/etc/glance/glance-api.conf

PATH=%%PREFIX%%/bin:$PATH

load_rc_config $name

: ${glance_api_enable:="NO"}

run_rc_command "$1"
