# misc/py-glance/files/openrc-glance-registry.in
#!/sbin/openrc-run

# glance_registry_logdir (str):
#	Set it to chagge log directory
#	Default is "/var/log/glance"
#
# glance_registry_args (str):
#	Set it to change command line arguments.
#	Default is "--log-file ${glance_registry_logdir}/glance-registry.log"
#

name=glance_registry

required_files=%%PREFIX%%/etc/glance/glance-registry.conf
PATH=%%PREFIX%%/bin:%%PREFIX%%/sbin:$PATH

pidfile="/var/run/glance-registry.pid"

: ${glance_registry_logdir:="/var/log/glance"}
: ${glance_registry_args:="--log-file ${glance_registry_logdir}/glance-registry.log"}

command="%%PREFIX%%/bin/python2.7"
supervisor=supervise-daemon
command_args="glance-registry ${glance_registry_args}"
command_user=glance

start_pre() {
	mkdir -p ${glance_registry_logdir}
	chown glance ${glance_registry_logdir}
}
depend()
{
	keyword -shutdown
}
