# dns/ddclient/files/ddclient.in
#!/bin/sh

# PROVIDE: ddclient
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable ddclient:
#
# ddclient_enable="YES"
#
# By default, ddclient does not daemonize.  To ensure ddclient runs
# as a daemon, the configuration file must contain a line like
# daemon=300 in it.
#

. /etc/rc.subr

name="ddclient"
rcvar=ddclient_enable

command="%%PREFIX%%/sbin/${name}"
required_files="%%ETCDIR%%/${name}.conf"

load_rc_config ${name}

: ${ddclient_enable="NO"}

run_rc_command "$1"
