# mail/scam-backscatter/files/openrc-scam-backscatter.in
#!/sbin/openrc-run

name="scambackscatter"
command="%%PREFIX%%/libexec/scam-back"

scambackscatter_user=${scambackscatter_user-"mailnull"}
scambackscatter_group=${scambackscatter_group-"mailnull"}

scambackscatter_dir=${scambackscatter_dir-"/var/scam-back"}
scambackscatter_runas=${scambackscatter_runas-"$scambackscatter_user"}
scambackscatter_pidfile=${scambackscatter_pidfile-"/var/run/scam-back.pid"}
scambackscatter_sockfile=${scambackscatter_sockfile-"$scambackscatter_dir/scam-back.sock"}
scambackscatter_cfgfile=${scambackscatter_cfgfile-"%%PREFIX%%/etc/mail/scam-back.conf"}
scambackscatter_flags=${scambackscatter_flags-"-p unix:$scambackscatter_sockfile \
-f $scambackscatter_cfgfile -u $scambackscatter_runas -P $scambackscatter_pidfile \
-b $scambackscatter_dir/backscatter.txt -D"}
command_user=mailnull

start_pre()
{
	if [ ! -d "${scambackscatter_dir}/." ]; then
		mkdir ${scambackscatter_dir}
		chown ${scambackscatter_user}:${scambackscatter_group} ${scambackscatter_dir}
		chmod 700 ${scambackscatter_dir}
	fi
	if [ ! -f "${scambackscatter_pidfile}" ]; then
		touch ${scambackscatter_pidfile}
		chown ${scambackscatter_user}:${scambackscatter_group} ${scambackscatter_pidfile}
	fi
}
