# sysutils/amazon-ssm-agent/files/openrc-amazon-ssm-agent.in
#!/sbin/openrc-run

name="amazon_ssm_agent"

pidfile="/var/run/${name}.pid"
command="%%PREFIX%%/sbin/amazon-ssm-agent"
supervisor=supervise-daemon

depend() {
	keyword -shutdown
}
