# palm/pilot-link/pkg-plist
etc/rc.d/pilot-link.sh
pilot/bin/addresses
pilot/bin/ccexample
pilot/bin/debugsh
pilot/bin/dlpsh
pilot/bin/getrom
pilot/bin/hinotes
pilot/bin/iambicexample
pilot/bin/ietf2datebook
pilot/bin/install-datebook
pilot/bin/install-hinote
pilot/bin/install-memo
pilot/bin/install-todos
pilot/bin/install-user
pilot/bin/memos
pilot/bin/money2qif
pilot/bin/pi-csd
pilot/bin/pi-getram
pilot/bin/pi-getrom
pilot/bin/pi-nredir
pilot/bin/pi-port
pilot/bin/pilot-addresses
pilot/bin/pilot-clip
pilot/bin/pilot-debug
pilot/bin/pilot-dedupe
pilot/bin/pilot-file
pilot/bin/pilot-mail
pilot/bin/pilot-prc
pilot/bin/pilot-schlep
pilot/bin/pilot-undelete
pilot/bin/pilot-xfer
pilot/bin/read-expenses
pilot/bin/read-ical
pilot/bin/read-todos
pilot/bin/reminders
pilot/bin/sync-plan
pilot/include/pi-address.h
pilot/include/pi-address.hxx
pilot/include/pi-appinfo.h
pilot/include/pi-appinfo.hxx
pilot/include/pi-args.h
pilot/include/pi-cmp.h
pilot/include/pi-config-win32.h
pilot/include/pi-config.h
pilot/include/pi-datebook.h
pilot/include/pi-datebook.hxx
pilot/include/pi-dlp.h
pilot/include/pi-dlp.hxx
pilot/include/pi-expense.h
pilot/include/pi-file.h
pilot/include/pi-header.h
pilot/include/pi-hinote.h
pilot/include/pi-iambicExpense.h
pilot/include/pi-iambicExpense.hxx
pilot/include/pi-inet.h
pilot/include/pi-inetserial.h
pilot/include/pi-macros.h
pilot/include/pi-mail.h
pilot/include/pi-memo.h
pilot/include/pi-memo.hxx
pilot/include/pi-money.h
pilot/include/pi-padp.h
pilot/include/pi-serial.h
pilot/include/pi-slp.h
pilot/include/pi-sockaddr-win32.h
pilot/include/pi-sockaddr.h
pilot/include/pi-socket.h
pilot/include/pi-source.h
pilot/include/pi-sync.h
pilot/include/pi-syspkt.h
pilot/include/pi-todo.h
pilot/include/pi-todo.hxx
pilot/include/pi-util.h
pilot/include/pi-version.h
pilot/lib/libpicc.a
pilot/lib/libpisock.a
pilot/lib/libpisock.so
pilot/lib/libpisock.so.4.0
pilot/lib/pix/b1.gif
pilot/lib/pix/b2.gif
pilot/lib/pix/b3.gif
pilot/lib/pix/b4.gif
pilot/lib/pix/b5.gif
pilot/lib/pix/b6.gif
pilot/lib/pix/b7.gif
pilot/lib/pix/case.gif
pilot/lib/pix/case.xbm
pilot/lib/pix/casemask.xbm
@dirrm pilot/lib/pix
@dirrm pilot/lib
@dirrm pilot/include
@dirrm pilot/bin
@dirrm pilot
