# sysutils/fusefs-kmod/pkg-plist
etc/rc.d/fuse.sh
modules/fuse.ko
sbin/mount_fusefs
sbin/mount_fusefs_safe
@unexec kldxref %D/modules
@unexec rmdir %D/modules 2>/dev/null || true
