# www/thttpd/pkg-plist
bin/htpasswd
bin/makeweb
bin/syslogtocern
etc/rc.d/thttpd.sh
etc/thttpd.conf.sample
sbin/thttpd
sbin/thttpd_wrapper
share/doc/thttpd/notes.html
www/cgi-bin/phf
www/cgi-bin/redirect
www/cgi-bin/ssi
@dirrm share/doc/thttpd
@dirrm www/cgi-bin
@dirrm www/data
@dirrm www/logs
@dirrm www
