# mail/raysfilter/pkg-plist
@comment $FreeBSD$
libexec/rays-filter
bin/mf-reset
bin/mf-start
bin/mf-stop
bin/mf-archive
bin/mf-analyze
bin/mf-memory
etc/rc.d/raysfilter.sh.sample
etc/mail-filter/header-list.conf-dist
etc/mail-filter/mail-filter.conf-dist
etc/mail-filter/string-list.conf-dist
@dirrm etc/mail-filter
%%PORTDOCS%%%%DOCSDIR%%/rays-filter.html
%%PORTDOCS%%%%DOCSDIR%%/install.html
%%PORTDOCS%%%%DOCSDIR%%/configure.html
%%PORTDOCS%%%%DOCSDIR%%/operate.html
%%PORTDOCS%%%%DOCSDIR%%/utilities.html
%%PORTDOCS%%%%DOCSDIR%%/GNU-public-licence.txt
%%PORTDOCS%%%%DOCSDIR%%/copyright.txt
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%@dirrm %%DOCSDIR%%
