# japanese/elvis/pkg-plist
bin/elvis
bin/ctags
bin/ref
bin/elvrec
bin/elvprsv
bin/vi
bin/ex
bin/view
bin/input
etc/rc.d/ja-elvis.sh
share/doc/ja-elvis/README.j18
@dirrm share/doc/ja-elvis
