# emulators/vmware-guestd/pkg-plist
sbin/checkvm
sbin/vmware-guestd
etc/rc.d/vmware-guestd.sh
