# sysutils/pfstat/files/openrc-pfstatd.in
#!/sbin/openrc-run

name="pfstatd"
command="%%PREFIX%%/bin/${name}"

depend() {
	keyword -shutdown
}
