# astro/squeuer/pkg-plist
sbin/squeuer
etc/squeuer.conf.sample
etc/rc.d/squeuer.sh.sample
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%@dirrm %%DOCSDIR%%
