# www/squid30/pkg-plist
etc/rc.d/squid.sh
etc/squid/errors/ERR_ACCESS_DENIED
etc/squid/errors/ERR_CACHE_ACCESS_DENIED
etc/squid/errors/ERR_CACHE_MGR_ACCESS_DENIED
etc/squid/errors/ERR_CANNOT_FORWARD
etc/squid/errors/ERR_CONNECT_FAIL
etc/squid/errors/ERR_DNS_FAIL
etc/squid/errors/ERR_FORWARDING_DENIED
etc/squid/errors/ERR_FTP_DISABLED
etc/squid/errors/ERR_FTP_FAILURE
etc/squid/errors/ERR_FTP_FORBIDDEN
etc/squid/errors/ERR_FTP_NOT_FOUND
etc/squid/errors/ERR_FTP_PUT_CREATED
etc/squid/errors/ERR_FTP_PUT_ERROR
etc/squid/errors/ERR_FTP_PUT_MODIFIED
etc/squid/errors/ERR_FTP_UNAVAILABLE
etc/squid/errors/ERR_INVALID_REQ
etc/squid/errors/ERR_INVALID_URL
etc/squid/errors/ERR_LIFETIME_EXP
etc/squid/errors/ERR_NO_RELAY
etc/squid/errors/ERR_ONLY_IF_CACHED_MISS
etc/squid/errors/ERR_READ_ERROR
etc/squid/errors/ERR_READ_TIMEOUT
etc/squid/errors/ERR_SHUTTING_DOWN
etc/squid/errors/ERR_SOCKET_FAILURE
etc/squid/errors/ERR_TOO_BIG
etc/squid/errors/ERR_UNSUP_REQ
etc/squid/errors/ERR_URN_RESOLVE
etc/squid/errors/ERR_WRITE_ERROR
etc/squid/errors/ERR_ZERO_SIZE_OBJECT
etc/squid/icons/anthony-binhex.gif
etc/squid/icons/anthony-bomb.gif
etc/squid/icons/anthony-box.gif
etc/squid/icons/anthony-box2.gif
etc/squid/icons/anthony-c.gif
etc/squid/icons/anthony-compressed.gif
etc/squid/icons/anthony-dir.gif
etc/squid/icons/anthony-dirup.gif
etc/squid/icons/anthony-dvi.gif
etc/squid/icons/anthony-f.gif
etc/squid/icons/anthony-image.gif
etc/squid/icons/anthony-image2.gif
etc/squid/icons/anthony-layout.gif
etc/squid/icons/anthony-link.gif
etc/squid/icons/anthony-movie.gif
etc/squid/icons/anthony-pdf.gif
etc/squid/icons/anthony-portal.gif
etc/squid/icons/anthony-ps.gif
etc/squid/icons/anthony-quill.gif
etc/squid/icons/anthony-script.gif
etc/squid/icons/anthony-sound.gif
etc/squid/icons/anthony-tar.gif
etc/squid/icons/anthony-tex.gif
etc/squid/icons/anthony-text.gif
etc/squid/icons/anthony-unknown.gif
etc/squid/icons/anthony-xbm.gif
etc/squid/icons/anthony-xpm.gif
etc/squid/mib.txt
etc/squid/mime.conf.default
etc/squid/squid.conf.default
libexec/cachemgr.cgi
libexec/dnsserver
libexec/pinger
libexec/unlinkd
sbin/RunAccel
sbin/RunCache
sbin/client
sbin/squid
@exec mkdir -p %D/squid/cache
@exec mkdir -p %D/squid/logs
@unexec rmdir %D/squid/cache 2>/dev/null || true
@unexec rmdir %D/squid/logs 2>/dev/null || true
@unexec rmdir %D/squid 2>/dev/null || true
