# net/openntpd/files/openrc-openntpd.in
#!/sbin/openrc-run
# Copyright (c) 2007-2015 The OpenRC Authors.
# See the Authors file at the top-level directory of this distribution and
# https://github.com/OpenRC/openrc/blob/master/AUTHORS
#
# This file is part of OpenRC. It is subject to the license terms in
# the LICENSE file found in the top-level directory of this
# distribution and at https://github.com/OpenRC/openrc/blob/master/LICENSE
# This file may not be copied, modified, propagated, or distributed
# except according to the terms contained in the LICENSE file.

: ${openntpd_config:=%%PREFIX%%/etc/ntpd.conf}
: ${ntpd_drift:=/var/db/ntpd.drift}


#Cleaned up by Ken Moore <ken@ixsystems.com> 6/20/17

name="openntpd"
command=%%PREFIX%%/sbin/ntpd
required_files=$openntpd_config
pidfile=/var/run/ntpd.pid
command_args="$ntpd_args -f $openntpd_config -d -s"
command_background="true"
description="Network Time Protocol Daemon"
start_stop_daemon_args="--stdout /var/log/openntpd.log --stderr /var/log/openntpd.log"

depend()
{
	need localmount
	provide ntpd
}
