# www/orion/pkg-plist
etc/rc.d/020.orion.sh
bin/orionctl
%%T%%/config/database-schemas/hypersonic.xml
%%T%%/config/database-schemas/ms-access.xml
%%T%%/config/database-schemas/ms-sql.xml
%%T%%/config/database-schemas/oracle.xml
%%T%%/config/database-schemas/postgresql.xml
%%T%%/config/database-schemas/sapdb.xml
%%T%%/config/database-schemas/sybase.xml
%%T%%/config/data-sources.xml
%%T%%/config/default-web-site.xml
%%T%%/config/global-web-application.xml
%%T%%/config/jms.xml
%%T%%/config/mime.types
%%T%%/config/principals.xml
%%T%%/config/rmi.xml
%%T%%/config/application.xml
%%T%%/config/server.xml
%%T%%/default-web-app/examples/cal/TableBean.txt
%%T%%/default-web-app/examples/cal/cal1.txt
%%T%%/default-web-app/examples/cal/cal2.jsp
%%T%%/default-web-app/examples/cal/cal2.txt
%%T%%/default-web-app/examples/cal/calendar.html
%%T%%/default-web-app/examples/cal/Entries.txt
%%T%%/default-web-app/examples/cal/Entry.txt
%%T%%/default-web-app/examples/cal/JspCalendar.txt
%%T%%/default-web-app/examples/cal/login.html
%%T%%/default-web-app/examples/cal/cal1.jsp
%%T%%/default-web-app/examples/images/code.gif
%%T%%/default-web-app/examples/images/execute.gif
%%T%%/default-web-app/examples/images/return.gif
%%T%%/default-web-app/examples/jsp/include/foo.html
%%T%%/default-web-app/examples/jsp/include/foo.jsp
%%T%%/default-web-app/examples/jsp/include/inc.html
%%T%%/default-web-app/examples/jsp/include/include.jsp
%%T%%/default-web-app/examples/jsp/include/include.txt
%%T%%/default-web-app/examples/jsp/checkbox/check.html
%%T%%/default-web-app/examples/jsp/checkbox/checkresult.jsp
%%T%%/default-web-app/examples/jsp/checkbox/checkresult.txt
%%T%%/default-web-app/examples/jsp/checkbox/CheckTest.html
%%T%%/default-web-app/examples/jsp/checkbox/cresult.html
%%T%%/default-web-app/examples/jsp/colors/clr.html
%%T%%/default-web-app/examples/jsp/colors/ColorGameBean.html
%%T%%/default-web-app/examples/jsp/colors/colors.html
%%T%%/default-web-app/examples/jsp/colors/colors.txt
%%T%%/default-web-app/examples/jsp/colors/colrs.jsp
%%T%%/default-web-app/examples/jsp/dates/date.html
%%T%%/default-web-app/examples/jsp/dates/date.jsp
%%T%%/default-web-app/examples/jsp/dates/date.txt
%%T%%/default-web-app/examples/jsp/error/er.html
%%T%%/default-web-app/examples/jsp/error/err.jsp
%%T%%/default-web-app/examples/jsp/error/err.txt
%%T%%/default-web-app/examples/jsp/error/error.html
%%T%%/default-web-app/examples/jsp/error/errorpge.jsp
%%T%%/default-web-app/examples/jsp/forward/forward.jsp
%%T%%/default-web-app/examples/jsp/forward/forward.txt
%%T%%/default-web-app/examples/jsp/forward/fwd.html
%%T%%/default-web-app/examples/jsp/forward/one.jsp
%%T%%/default-web-app/examples/jsp/forward/two.html
%%T%%/default-web-app/examples/jsp/snp/snoop.html
%%T%%/default-web-app/examples/jsp/snp/snoop.jsp
%%T%%/default-web-app/examples/jsp/snp/snoop.txt
%%T%%/default-web-app/examples/jsp/cal/TableBean.txt
%%T%%/default-web-app/examples/jsp/cal/cal1.txt
%%T%%/default-web-app/examples/jsp/cal/cal2.jsp
%%T%%/default-web-app/examples/jsp/cal/cal2.txt
%%T%%/default-web-app/examples/jsp/cal/calendar.html
%%T%%/default-web-app/examples/jsp/cal/Entries.txt
%%T%%/default-web-app/examples/jsp/cal/Entry.txt
%%T%%/default-web-app/examples/jsp/cal/JspCalendar.txt
%%T%%/default-web-app/examples/jsp/cal/login.html
%%T%%/default-web-app/examples/jsp/cal/cal1.jsp
%%T%%/default-web-app/examples/jsp/jsptoserv/hello.jsp
%%T%%/default-web-app/examples/jsp/jsptoserv/jsptoservlet.jsp
%%T%%/default-web-app/examples/jsp/jsptoserv/jts.html
%%T%%/default-web-app/examples/jsp/jsptoserv/jts.txt
%%T%%/default-web-app/examples/jsp/jsptoserv/stj.txt
%%T%%/default-web-app/examples/jsp/num/numguess.html
%%T%%/default-web-app/examples/jsp/num/numguess.jsp
%%T%%/default-web-app/examples/jsp/num/numguess.txt
%%T%%/default-web-app/examples/jsp/plugin/applet/Clock2.class
%%T%%/default-web-app/examples/jsp/plugin/applet/Clock2.java
%%T%%/default-web-app/examples/jsp/plugin/plugin.html
%%T%%/default-web-app/examples/jsp/plugin/plugin.jsp
%%T%%/default-web-app/examples/jsp/plugin/plugin.txt
%%T%%/default-web-app/examples/jsp/sessions/carts.html
%%T%%/default-web-app/examples/jsp/sessions/carts.jsp
%%T%%/default-web-app/examples/jsp/sessions/carts.txt
%%T%%/default-web-app/examples/jsp/sessions/crt.html
%%T%%/default-web-app/examples/jsp/sessions/DummyCart.html
%%T%%/default-web-app/examples/jsp/taglib/loop/META-INF/taglib.tld
%%T%%/default-web-app/examples/jsp/taglib/loop/taglib/LoopTag.class
%%T%%/default-web-app/examples/jsp/taglib/loop/taglib/LoopTag.java
%%T%%/default-web-app/examples/jsp/taglib/loop/index.html
%%T%%/default-web-app/examples/jsp/taglib/loop/looptag.jsp
%%T%%/default-web-app/examples/jsp/taglib/loop/looptag.jsp.txt
%%T%%/default-web-app/examples/jsp/taglib/loop/looptags.jar
%%T%%/default-web-app/examples/jsp/index.html
%%T%%/default-web-app/examples/perl/helloworld.pl
%%T%%/default-web-app/examples/perl/post.html
%%T%%/default-web-app/examples/perl/post.pl
%%T%%/default-web-app/examples/ssi/other.txt
%%T%%/default-web-app/examples/ssi/real.txt
%%T%%/default-web-app/examples/ssi/ssi.shtml
%%T%%/default-web-app/examples/xsl/data.xml
%%T%%/default-web-app/examples/xsl/data.xsl
%%T%%/default-web-app/examples/xsl/doc.dtd
%%T%%/default-web-app/examples/xsl/doc.jsp
%%T%%/default-web-app/examples/xsl/doc.xml
%%T%%/default-web-app/examples/xsl/doc.xsl
%%T%%/default-web-app/WEB-INF/classes/colors/ColorGameBean.class
%%T%%/default-web-app/WEB-INF/classes/colors/ColorGameBean.java
%%T%%/default-web-app/WEB-INF/classes/num/NumberGuessBean.class
%%T%%/default-web-app/WEB-INF/classes/num/NumberGuessBean.java
%%T%%/default-web-app/WEB-INF/classes/error/Smart.class
%%T%%/default-web-app/WEB-INF/classes/error/Smart.java
%%T%%/default-web-app/WEB-INF/classes/sessions/DummyCart.class
%%T%%/default-web-app/WEB-INF/classes/sessions/DummyCart.java
%%T%%/default-web-app/WEB-INF/classes/dates/JspCalendar.class
%%T%%/default-web-app/WEB-INF/classes/dates/JspCalendar.java
%%T%%/default-web-app/WEB-INF/classes/cal/TableBean.java
%%T%%/default-web-app/WEB-INF/classes/cal/Entries.java
%%T%%/default-web-app/WEB-INF/classes/cal/Entry.class
%%T%%/default-web-app/WEB-INF/classes/cal/Entry.java
%%T%%/default-web-app/WEB-INF/classes/cal/JspCalendar.class
%%T%%/default-web-app/WEB-INF/classes/cal/JspCalendar.java
%%T%%/default-web-app/WEB-INF/classes/cal/TableBean.class
%%T%%/default-web-app/WEB-INF/classes/cal/Entries.class
%%T%%/default-web-app/WEB-INF/classes/checkbox/CheckTest.class
%%T%%/default-web-app/WEB-INF/classes/checkbox/CheckTest.java
%%T%%/default-web-app/WEB-INF/classes/Counter.java
%%T%%/default-web-app/WEB-INF/classes/DateServlet.class
%%T%%/default-web-app/WEB-INF/classes/DateServlet.java
%%T%%/default-web-app/WEB-INF/classes/FingerServlet.class
%%T%%/default-web-app/WEB-INF/classes/FingerServlet.java
%%T%%/default-web-app/WEB-INF/classes/HelloWorldExample.class
%%T%%/default-web-app/WEB-INF/classes/HelloWorldExample.java
%%T%%/default-web-app/WEB-INF/classes/HelloWorldServlet.class
%%T%%/default-web-app/WEB-INF/classes/HelloWorldServlet.java
%%T%%/default-web-app/WEB-INF/classes/LinkCheckerServlet.class
%%T%%/default-web-app/WEB-INF/classes/LinkCheckerServlet.java
%%T%%/default-web-app/WEB-INF/classes/LocalStrings.properties
%%T%%/default-web-app/WEB-INF/classes/ParameterBean.class
%%T%%/default-web-app/WEB-INF/classes/ParameterBean.java
%%T%%/default-web-app/WEB-INF/classes/RedirectServlet.class
%%T%%/default-web-app/WEB-INF/classes/RedirectServlet.java
%%T%%/default-web-app/WEB-INF/classes/RequestHeaderExample.class
%%T%%/default-web-app/WEB-INF/classes/RequestHeaderExample.java
%%T%%/default-web-app/WEB-INF/classes/RequestInfoExample.class
%%T%%/default-web-app/WEB-INF/classes/RequestInfoExample.java
%%T%%/default-web-app/WEB-INF/classes/RequestParamExample.class
%%T%%/default-web-app/WEB-INF/classes/RequestParamExample.java
%%T%%/default-web-app/WEB-INF/classes/SessionExample.class
%%T%%/default-web-app/WEB-INF/classes/SessionExample.java
%%T%%/default-web-app/WEB-INF/classes/SessionServlet.class
%%T%%/default-web-app/WEB-INF/classes/SessionServlet.java
%%T%%/default-web-app/WEB-INF/classes/SimpleServlet.class
%%T%%/default-web-app/WEB-INF/classes/SimpleServlet.java
%%T%%/default-web-app/WEB-INF/classes/SnoopServlet.class
%%T%%/default-web-app/WEB-INF/classes/SnoopServlet.java
%%T%%/default-web-app/WEB-INF/classes/SurveyServlet.class
%%T%%/default-web-app/WEB-INF/classes/SurveyServlet.java
%%T%%/default-web-app/WEB-INF/classes/ToJSPServlet.class
%%T%%/default-web-app/WEB-INF/classes/ToJSPServlet.java
%%T%%/default-web-app/WEB-INF/classes/UpperCaseFilter.class
%%T%%/default-web-app/WEB-INF/classes/UpperCaseFilter.java
%%T%%/default-web-app/WEB-INF/classes/BBoardServlet.class
%%T%%/default-web-app/WEB-INF/classes/URLInfo.class
%%T%%/default-web-app/WEB-INF/classes/BBoardServlet.java
%%T%%/default-web-app/WEB-INF/classes/ConfigServlet.class
%%T%%/default-web-app/WEB-INF/classes/ConfigServlet.java
%%T%%/default-web-app/WEB-INF/classes/CookieExample.class
%%T%%/default-web-app/WEB-INF/classes/CookieExample.java
%%T%%/default-web-app/WEB-INF/classes/Counter.class
%%T%%/default-web-app/WEB-INF/web.xml
%%T%%/default-web-app/index.html
%%T%%/applications/news.ear
%%T%%/applications/atm.ear
%%T%%/applications/estore-bugs.txt
%%T%%/applications/estore-install.txt
%%T%%/applications/estore-orionadapter.jar
%%T%%/applications/estore-patched.ear
%%T%%/applications/news-install.txt
%%T%%/applications/atm-install.txt
%%T%%/demo/ejb/cart/META-INF/application-client.xml
%%T%%/demo/ejb/cart/META-INF/ejb-jar.xml
%%T%%/demo/ejb/cart/META-INF/MANIFEST.MF
%%T%%/demo/ejb/cart/Cart.java
%%T%%/demo/ejb/cart/CartClient.class
%%T%%/demo/ejb/cart/CartClient.java
%%T%%/demo/ejb/cart/CartEJB.class
%%T%%/demo/ejb/cart/CartEJB.java
%%T%%/demo/ejb/cart/CartHome.class
%%T%%/demo/ejb/cart/CartHome.java
%%T%%/demo/ejb/cart/jndi.properties
%%T%%/demo/ejb/cart/NotInCartException.class
%%T%%/demo/ejb/cart/NotInCartException.java
%%T%%/demo/ejb/cart/Cart.class
%%T%%/demo/ejb/cart/orion-cart.jar
%%T%%/demo/ejb/META-INF/application.xml
%%T%%/demo/ejb/product/META-INF/application-client.xml
%%T%%/demo/ejb/product/META-INF/ejb-jar.xml
%%T%%/demo/ejb/product/META-INF/MANIFEST.MF
%%T%%/demo/ejb/product/orion-product.jar
%%T%%/demo/ejb/product/Product.class
%%T%%/demo/ejb/product/Product.java
%%T%%/demo/ejb/product/product.skeleton
%%T%%/demo/ejb/product/ProductClient.class
%%T%%/demo/ejb/product/ProductClient.java
%%T%%/demo/ejb/product/ProductEJB.class
%%T%%/demo/ejb/product/ProductEJB.java
%%T%%/demo/ejb/product/ProductHome.class
%%T%%/demo/ejb/product/ProductHome.java
%%T%%/demo/ejb/product/ProductServlet.class
%%T%%/demo/ejb/product/jndi.properties
%%T%%/demo/ejb/product/ProductServlet.java
%%T%%/demo/ejb/usermanager/META-INF/ejb-jar.xml
%%T%%/demo/ejb/install.txt
%%T%%/demo/jms/coffeemaker/META-INF/application-client.xml
%%T%%/demo/jms/coffeemaker/META-INF/MANIFEST.MF
%%T%%/demo/jms/coffeemaker/META-INF/orion-application-client.xml
%%T%%/demo/jms/coffeemaker/CoffeeMachine.class
%%T%%/demo/jms/coffeemaker/CoffeeMachine.java
%%T%%/demo/jms/coffeemaker/CoffeeMaker.class
%%T%%/demo/jms/coffeemaker/CoffeeMaker.java
%%T%%/demo/jms/coffeemaker/jndi.properties
%%T%%/demo/jms/jmschat/META-INF/application-client.xml
%%T%%/demo/jms/jmschat/META-INF/MANIFEST.MF
%%T%%/demo/jms/jmschat/META-INF/orion-application-client.xml
%%T%%/demo/jms/jmschat/JMSChat$1.class
%%T%%/demo/jms/jmschat/JMSChat.class
%%T%%/demo/jms/jmschat/JMSChat.java
%%T%%/demo/jms/jmschat/jndi.properties
%%T%%/demo/jms/coffeemaker.jar
%%T%%/demo/jms/jmschat.jar
%%T%%/demo/jms/jmsdemo-readme.txt
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsView.java
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsItem.java
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsItemClient.class
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsItemEJB.class
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsItemEJB.java
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsItemHome.class
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsItemHome.java
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsItemServlet.class
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsItemServlet.java
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsSession.class
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsSession.java
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsSessionEJB.class
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsSessionEJB.java
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsSessionHome.class
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsSessionHome.java
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsView.class
%%T%%/demo/news-app/news-ejb-source/com/evermind/ejb/NewsItem.class
%%T%%/demo/news-app/news-ejb-source/META-INF/ejb-jar.xml
%%T%%/demo/news-app/news-client-source/com/evermind/ejb/NewsClient.class
%%T%%/demo/news-app/news-client-source/com/evermind/ejb/NewsClient.java
%%T%%/demo/news-app/news-client-source/com/evermind/ejb/NewsItemClient.class
%%T%%/demo/news-app/news-client-source/META-INF/application-client.xml
%%T%%/demo/news-app/news-client-source/META-INF/MANIFEST.MF
%%T%%/demo/news-app/news-web/admin/editnews.jsp
%%T%%/demo/news-app/news-web/admin/editnews_done.jsp
%%T%%/demo/news-app/news-web/WEB-INF/lib/ejbtags.jar
%%T%%/demo/news-app/news-web/WEB-INF/lib/utiltags.jar
%%T%%/demo/news-app/news-web/WEB-INF/web.xml
%%T%%/demo/news-app/news-web/border.html
%%T%%/demo/news-app/news-web/checks.jsp
%%T%%/demo/news-app/news-web/index.html
%%T%%/demo/news-app/news-web/initviews.jsp
%%T%%/demo/news-app/news-web/login-error.html
%%T%%/demo/news-app/news-web/addnews.jsp
%%T%%/demo/news-app/news-web/loginbox.jsp
%%T%%/demo/news-app/news-web/logout.jsp
%%T%%/demo/news-app/news-web/news-body.jsp
%%T%%/demo/news-app/news-web/news.jsp
%%T%%/demo/news-app/news-web/newsdetail.jsp
%%T%%/demo/news-app/news-web/newssearch.jsp
%%T%%/demo/news-app/news-web/usertest.jsp
%%T%%/demo/news-app/news-web/login.html
%%T%%/demo/news-app/news-web/addnews_done.jsp
%%T%%/demo/news-app/META-INF/application.xml
%%T%%/demo/news-app/news-web.war
%%T%%/demo/news-app/news-client.jar
%%T%%/demo/news-app/news-ejb.jar
%%T%%/demo/news-app/counter.jar
%%T%%/demo/rmi/RMIApp.class
%%T%%/demo/rmi/RMIApp.java
%%T%%/demo/ssl/ssl-user-registration.jsp
%%T%%/demo/messagelogger.ear
%%T%%/lib/hsql.jar
%%T%%/log/stdout.log
%%T%%/log/stderr.log
%%T%%/loadbalancer.jar
%%T%%/jta.jar
%%T%%/jsse.jar
%%T%%/jnet.jar
%%T%%/jndi.jar
%%T%%/jdbc.jar
%%T%%/jcert.jar
%%T%%/jaxp.jar
%%T%%/ejbmaker.jar
%%T%%/ejbassembler.jar
%%T%%/ejb.jar
%%T%%/earassembler.jar
%%T%%/clientassembler.jar
%%T%%/changes.txt
%%T%%/autoupdate.properties
%%T%%/autoupdate.jar
%%T%%/assemblerlauncher.jar
%%T%%/applicationlauncher.jar
%%T%%/admin.jar
%%T%%/activation.jar
%%T%%/Readme.txt
%%T%%/xerces.jar
%%T%%/Copyright.txt
%%T%%/xalan.jar
%%T%%/webappassembler.jar
%%T%%/taglibassembler.jar
%%T%%/parser.jar
%%T%%/orionconsole.jar
%%T%%/orion.properties
%%T%%/orion.jar
%%T%%/tools.jar
%%T%%/mail.jar
%%T%%/crimson.jar
%%T%%/jaas.jar
@dirrm %%T%%/application-deployments
@dirrm %%T%%/config/database-schemas
@dirrm %%T%%/config
@dirrm %%T%%/default-web-app/examples/cal
@dirrm %%T%%/default-web-app/examples/images
@dirrm %%T%%/default-web-app/examples/jsp/include
@dirrm %%T%%/default-web-app/examples/jsp/checkbox
@dirrm %%T%%/default-web-app/examples/jsp/colors
@dirrm %%T%%/default-web-app/examples/jsp/dates
@dirrm %%T%%/default-web-app/examples/jsp/error
@dirrm %%T%%/default-web-app/examples/jsp/forward
@dirrm %%T%%/default-web-app/examples/jsp/snp
@dirrm %%T%%/default-web-app/examples/jsp/cal
@dirrm %%T%%/default-web-app/examples/jsp/jsptoserv
@dirrm %%T%%/default-web-app/examples/jsp/num
@dirrm %%T%%/default-web-app/examples/jsp/plugin/applet
@dirrm %%T%%/default-web-app/examples/jsp/plugin
@dirrm %%T%%/default-web-app/examples/jsp/sessions
@dirrm %%T%%/default-web-app/examples/jsp/taglib/loop/META-INF
@dirrm %%T%%/default-web-app/examples/jsp/taglib/loop/taglib
@dirrm %%T%%/default-web-app/examples/jsp/taglib/loop
@dirrm %%T%%/default-web-app/examples/jsp/taglib
@dirrm %%T%%/default-web-app/examples/jsp
@dirrm %%T%%/default-web-app/examples/perl
@dirrm %%T%%/default-web-app/examples/ssi
@dirrm %%T%%/default-web-app/examples/xsl
@dirrm %%T%%/default-web-app/examples
@dirrm %%T%%/default-web-app/WEB-INF/classes/colors
@dirrm %%T%%/default-web-app/WEB-INF/classes/num
@dirrm %%T%%/default-web-app/WEB-INF/classes/error
@dirrm %%T%%/default-web-app/WEB-INF/classes/sessions
@dirrm %%T%%/default-web-app/WEB-INF/classes/dates
@dirrm %%T%%/default-web-app/WEB-INF/classes/cal
@dirrm %%T%%/default-web-app/WEB-INF/classes/checkbox
@dirrm %%T%%/default-web-app/WEB-INF/classes
@dirrm %%T%%/default-web-app/WEB-INF
@dirrm %%T%%/default-web-app
@dirrm %%T%%/applications
@dirrm %%T%%/database
@dirrm %%T%%/demo/ejb/cart/META-INF
@dirrm %%T%%/demo/ejb/cart
@dirrm %%T%%/demo/ejb/META-INF
@dirrm %%T%%/demo/ejb/product/META-INF
@dirrm %%T%%/demo/ejb/product
@dirrm %%T%%/demo/ejb/usermanager/META-INF
@dirrm %%T%%/demo/ejb/usermanager
@dirrm %%T%%/demo/ejb
@dirrm %%T%%/demo/jms/coffeemaker/META-INF
@dirrm %%T%%/demo/jms/coffeemaker
@dirrm %%T%%/demo/jms/jmschat/META-INF
@dirrm %%T%%/demo/jms/jmschat
@dirrm %%T%%/demo/jms
@dirrm %%T%%/demo/news-app/news-ejb-source/com/evermind/ejb
@dirrm %%T%%/demo/news-app/news-ejb-source/com/evermind
@dirrm %%T%%/demo/news-app/news-ejb-source/com
@dirrm %%T%%/demo/news-app/news-ejb-source/META-INF
@dirrm %%T%%/demo/news-app/news-ejb-source
@dirrm %%T%%/demo/news-app/news-client-source/com/evermind/ejb
@dirrm %%T%%/demo/news-app/news-client-source/com/evermind
@dirrm %%T%%/demo/news-app/news-client-source/com
@dirrm %%T%%/demo/news-app/news-client-source/META-INF
@dirrm %%T%%/demo/news-app/news-client-source
@dirrm %%T%%/demo/news-app/news-web/admin
@dirrm %%T%%/demo/news-app/news-web/WEB-INF/classes
@dirrm %%T%%/demo/news-app/news-web/WEB-INF/lib
@dirrm %%T%%/demo/news-app/news-web/WEB-INF
@dirrm %%T%%/demo/news-app/news-web
@dirrm %%T%%/demo/news-app/META-INF
@dirrm %%T%%/demo/news-app
@dirrm %%T%%/demo/rmi
@dirrm %%T%%/demo/ssl
@dirrm %%T%%/demo
@dirrm %%T%%/lib
@dirrm %%T%%/log
@dirrm %%T%%/persistence/ejb
@dirrm %%T%%/persistence
@dirrm %%T%%
