# sysutils/qchroot/pkg-plist
bin/qchroot
etc/rc.d/qchroot.bootime
share/man/man8/qchroot.8.gz
