# www/butterfly/files/openrc-butterfly.in
#!/sbin/openrc-run

: ${butterfly_user="root"}
: ${butterfly_args="--unsecure"}

name="butterfly"
descrption=""
pidfile=/var/run/${name}.pid

command="%%PYTHON_CMD%%"

start() {
  echo "Starting ${name}."
  if ! [ -r %%ETCDIR%%/ssl/butterfly_localhost.crt -a %%ETCDIR%%/ssl/butterfly_localhost.key -a %%ETCDIR%%/ssl/butterfly_ca.crt ]; then
    %%PREFIX%%/bin/butterfly.server.py --generate-certs --host=localhost &&
    echo "${name}: generated the ssl certificate in %%ETCDIR%%/ssl"
  fi
  daemon -p $pidfile -u ${butterfly_user} %%PREFIX%%/bin/butterfly.server.py ${butterfly_args} >/var/log/${name}.log 2>&1
}
