# www/havp/pkg-plist
sbin/havp
@sample etc/havp/havp.config.default etc/havp/havp.config
@unexec [ ! -f %%RUN_DIR%%/havp.pid ] || %D/etc/rc.d/havp stop
%%EXAMPLESDIR%%/blacklist.sample
%%EXAMPLESDIR%%/templates/br/blacklist.html
%%EXAMPLESDIR%%/templates/br/dns.html
%%EXAMPLESDIR%%/templates/br/down.html
%%EXAMPLESDIR%%/templates/br/error.html
%%EXAMPLESDIR%%/templates/br/invalid.html
%%EXAMPLESDIR%%/templates/br/maxsize.html
%%EXAMPLESDIR%%/templates/br/request.html
%%EXAMPLESDIR%%/templates/br/scanner.html
%%EXAMPLESDIR%%/templates/br/virus.html
%%EXAMPLESDIR%%/templates/css2/havp.css
%%EXAMPLESDIR%%/templates/de/blacklist.html
%%EXAMPLESDIR%%/templates/de/dns.html
%%EXAMPLESDIR%%/templates/de/down.html
%%EXAMPLESDIR%%/templates/de/error.html
%%EXAMPLESDIR%%/templates/de/invalid.html
%%EXAMPLESDIR%%/templates/de/maxsize.html
%%EXAMPLESDIR%%/templates/de/request.html
%%EXAMPLESDIR%%/templates/de/scanner.html
%%EXAMPLESDIR%%/templates/de/virus.html
%%EXAMPLESDIR%%/templates/en/blacklist.html
%%EXAMPLESDIR%%/templates/en/dns.html
%%EXAMPLESDIR%%/templates/en/down.html
%%EXAMPLESDIR%%/templates/en/error.html
%%EXAMPLESDIR%%/templates/en/invalid.html
%%EXAMPLESDIR%%/templates/en/maxsize.html
%%EXAMPLESDIR%%/templates/en/request.html
%%EXAMPLESDIR%%/templates/en/scanner.html
%%EXAMPLESDIR%%/templates/en/virus.html
%%EXAMPLESDIR%%/templates/es/blacklist.html
%%EXAMPLESDIR%%/templates/es/dns.html
%%EXAMPLESDIR%%/templates/es/down.html
%%EXAMPLESDIR%%/templates/es/error.html
%%EXAMPLESDIR%%/templates/es/invalid.html
%%EXAMPLESDIR%%/templates/es/maxsize.html
%%EXAMPLESDIR%%/templates/es/request.html
%%EXAMPLESDIR%%/templates/es/scanner.html
%%EXAMPLESDIR%%/templates/es/virus.html
%%EXAMPLESDIR%%/templates/fr/blacklist.html
%%EXAMPLESDIR%%/templates/fr/dns.html
%%EXAMPLESDIR%%/templates/fr/down.html
%%EXAMPLESDIR%%/templates/fr/error.html
%%EXAMPLESDIR%%/templates/fr/invalid.html
%%EXAMPLESDIR%%/templates/fr/maxsize.html
%%EXAMPLESDIR%%/templates/fr/request.html
%%EXAMPLESDIR%%/templates/fr/scanner.html
%%EXAMPLESDIR%%/templates/fr/virus.html
%%EXAMPLESDIR%%/templates/it/blacklist.html
%%EXAMPLESDIR%%/templates/it/dns.html
%%EXAMPLESDIR%%/templates/it/down.html
%%EXAMPLESDIR%%/templates/it/error.html
%%EXAMPLESDIR%%/templates/it/invalid.html
%%EXAMPLESDIR%%/templates/it/maxsize.html
%%EXAMPLESDIR%%/templates/it/request.html
%%EXAMPLESDIR%%/templates/it/scanner.html
%%EXAMPLESDIR%%/templates/it/virus.html
%%EXAMPLESDIR%%/templates/nl/blacklist.html
%%EXAMPLESDIR%%/templates/nl/dns.html
%%EXAMPLESDIR%%/templates/nl/down.html
%%EXAMPLESDIR%%/templates/nl/error.html
%%EXAMPLESDIR%%/templates/nl/invalid.html
%%EXAMPLESDIR%%/templates/nl/maxsize.html
%%EXAMPLESDIR%%/templates/nl/request.html
%%EXAMPLESDIR%%/templates/nl/scanner.html
%%EXAMPLESDIR%%/templates/nl/virus.html
%%EXAMPLESDIR%%/templates/pf/blacklist.html
%%EXAMPLESDIR%%/templates/pf/dns.html
%%EXAMPLESDIR%%/templates/pf/down.html
%%EXAMPLESDIR%%/templates/pf/error.html
%%EXAMPLESDIR%%/templates/pf/invalid.html
%%EXAMPLESDIR%%/templates/pf/maxsize.html
%%EXAMPLESDIR%%/templates/pf/request.html
%%EXAMPLESDIR%%/templates/pf/scanner.html
%%EXAMPLESDIR%%/templates/pf/virus.html
%%EXAMPLESDIR%%/templates/pl/blacklist.html
%%EXAMPLESDIR%%/templates/pl/dns.html
%%EXAMPLESDIR%%/templates/pl/down.html
%%EXAMPLESDIR%%/templates/pl/error.html
%%EXAMPLESDIR%%/templates/pl/invalid.html
%%EXAMPLESDIR%%/templates/pl/maxsize.html
%%EXAMPLESDIR%%/templates/pl/request.html
%%EXAMPLESDIR%%/templates/pl/scanner.html
%%EXAMPLESDIR%%/templates/pl/virus.html
%%EXAMPLESDIR%%/templates/ru/blacklist.html
%%EXAMPLESDIR%%/templates/ru/dns.html
%%EXAMPLESDIR%%/templates/ru/down.html
%%EXAMPLESDIR%%/templates/ru/error.html
%%EXAMPLESDIR%%/templates/ru/invalid.html
%%EXAMPLESDIR%%/templates/ru/maxsize.html
%%EXAMPLESDIR%%/templates/ru/request.html
%%EXAMPLESDIR%%/templates/ru/scanner.html
%%EXAMPLESDIR%%/templates/ru/virus.html
%%EXAMPLESDIR%%/templates/sv/blacklist.html
%%EXAMPLESDIR%%/templates/sv/dns.html
%%EXAMPLESDIR%%/templates/sv/down.html
%%EXAMPLESDIR%%/templates/sv/error.html
%%EXAMPLESDIR%%/templates/sv/invalid.html
%%EXAMPLESDIR%%/templates/sv/maxsize.html
%%EXAMPLESDIR%%/templates/sv/request.html
%%EXAMPLESDIR%%/templates/sv/scanner.html
%%EXAMPLESDIR%%/templates/sv/virus.html
%%EXAMPLESDIR%%/whitelist.sample
@dir(havp,havp,755) %%LOG_DIR%%
@dir(havp,havp,755) %%TMP_DIR%%
@dir(havp,havp,755) %%RUN_DIR%%
