# mail/sid-milter/pkg-plist
etc/rc.d/milter-sid.sh.sample
libexec/sid-filter
%%WITH_AR%%include/libar.h
%%WITH_AR%%lib/libar.a
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/KNOWNBUGS
%%PORTDOCS%%%%DOCSDIR%%/LICENSE
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README-SenderID
%%PORTDOCS%%%%DOCSDIR%%/RELEASE_NOTES
%%PORTDOCS%%@dirrm %%DOCSDIR%%
