# multimedia/mjpg-streamer/files/mjpg_streamer.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: mjpg_streamer
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# mjpg_streamer_enable="YES"
# #optional
# mjpg_streamer_flags=""

. /etc/rc.subr

name=mjpg_streamer
rcvar=mjpg_streamer_enable

command=%%PREFIX%%/bin/${name}
command_args="-b"

mjpg_streamer_enable="NO"
mjpg_streamer_flags="-o 'output_http.so -w %%WWWDIR%%'"

load_rc_config $name
run_rc_command "$1"
