# multimedia/mjpg-streamer/files/openrc-mjpg_streamer.in
#!/sbin/openrc-run

name=mjpg_streamer

command=%%PREFIX%%/bin/${name}
command_args="-b -o 'output_http.so -w %%WWWDIR%%'"

depend()
{
	keyword -shutdown
}
