# net-mgmt/bandwidthd/pkg-plist
bandwidthd/bandwidthd
@sample bandwidthd/etc/bandwidthd.conf.sample
bandwidthd/htdocs/legend.gif
bandwidthd/htdocs/logo.gif
etc/rc.d/bandwidthd
%%PORTDOCS%%%%DOCSDIR%%/README
