# net/google-startup-scripts/files/google-startup.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: google-startup
# BEFORE: LOGIN
# REQUIRE: syslogd
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable google-startup:
#
# google_startup_enable="YES"
#

. /etc/rc.subr

name="google_startup"
rcvar="google_startup_enable"
start_cmd="google_start"
stop_cmd="google_stop"

google_start()
{
  %%PREFIX%%/share/google/onboot
  %%PREFIX%%/share/google/run-startup-scripts
}

google_stop()
{
  %%PREFIX%%/share/google/run-shutdown-scripts
}

google_startup_enable=${google_startup_enable:-"NO"}

load_rc_config "${name}"
run_rc_command "$1"
