# ftp/ftp4all/pkg-plist
sbin/f4a.ftpd
sbin/ftps
sbin/ftpa
etc/rc.d/ftp4all.sh.sample
etc/ftp4all/cdpath.cfg.sample
etc/ftp4all/check.cfg.sample
etc/ftp4all/checkdupe.cfg.sample
etc/ftp4all/dirshortcut.cfg.sample
etc/ftp4all/ftpd.reg.sample
etc/ftp4all/limits.cfg.sample
etc/ftp4all/msg/goodbye.sample
etc/ftp4all/msg/startup.sample
etc/ftp4all/msg/welcome.sample
etc/ftp4all/msg/who_body.sample
etc/ftp4all/msg/who_tail.sample
etc/ftp4all/msg/who_top.sample
share/doc/ftp4all/COPYRIGHT
share/doc/ftp4all/README
share/doc/ftp4all/RELEASE-TEXT
share/doc/ftp4all/VERSION
share/doc/ftp4all/CHANGES
share/doc/ftp4all/ftp4all.gif
share/doc/ftp4all/f4adp/wa_del.txt
share/doc/ftp4all/f4adp/aen.txt
share/doc/ftp4all/f4adp/cfg_avl.txt
share/doc/ftp4all/f4adp/cfg_get.txt
share/doc/ftp4all/f4adp/cfg_rm.txt
share/doc/ftp4all/f4adp/cfg_set.txt
share/doc/ftp4all/f4adp/cnt.txt
share/doc/ftp4all/f4adp/col.txt
share/doc/ftp4all/f4adp/dgf.txt
share/doc/ftp4all/f4adp/duf.txt
share/doc/ftp4all/f4adp/grp.txt
share/doc/ftp4all/f4adp/grp_add.txt
share/doc/ftp4all/f4adp/grp_chg.txt
share/doc/ftp4all/f4adp/grp_get.txt
share/doc/ftp4all/f4adp/grp_ls.txt
share/doc/ftp4all/f4adp/grp_lsr.txt
share/doc/ftp4all/f4adp/grp_mdn.txt
share/doc/ftp4all/f4adp/grp_mup.txt
share/doc/ftp4all/f4adp/grp_rm.txt
share/doc/ftp4all/f4adp/grp_set.txt
share/doc/ftp4all/f4adp/msg_avl.txt
share/doc/ftp4all/f4adp/msg_dfl.txt
share/doc/ftp4all/f4adp/msg_get.txt
share/doc/ftp4all/f4adp/msg_rm.txt
share/doc/ftp4all/f4adp/msg_set.txt
share/doc/ftp4all/f4adp/ulp.txt
share/doc/ftp4all/f4adp/prt_get.txt
share/doc/ftp4all/f4adp/prt_set.txt
share/doc/ftp4all/f4adp/srv_gst.txt
share/doc/ftp4all/f4adp/srv_sst.txt
share/doc/ftp4all/f4adp/f4adp_frame.html
share/doc/ftp4all/f4adp/f4adp_index.html
share/doc/ftp4all/f4adp/f4adp_overview.html
share/doc/ftp4all/f4adp/dup.txt
share/doc/ftp4all/f4adp/fcm.txt
share/doc/ftp4all/f4adp/fco.txt
share/doc/ftp4all/f4adp/fdn.txt
share/doc/ftp4all/f4adp/fup.txt
share/doc/ftp4all/f4adp/frm.txt
share/doc/ftp4all/f4adp/frn.txt
share/doc/ftp4all/f4adp/lim_gbw.txt
share/doc/ftp4all/f4adp/lim_rbw.txt
share/doc/ftp4all/f4adp/log_act.txt
share/doc/ftp4all/f4adp/log_beg.txt
share/doc/ftp4all/f4adp/log_end.txt
share/doc/ftp4all/f4adp/_template.txt
share/doc/ftp4all/f4adp/quit.txt
share/doc/ftp4all/f4adp/reg_get.txt
share/doc/ftp4all/f4adp/reg_key.txt
share/doc/ftp4all/f4adp/reg_ls.txt
share/doc/ftp4all/f4adp/uan.txt
share/doc/ftp4all/f4adp/reg_rm.txt
share/doc/ftp4all/f4adp/reg_set.txt
share/doc/ftp4all/f4adp/usr_add.txt
share/doc/ftp4all/f4adp/usr_adq.txt
share/doc/ftp4all/f4adp/usr_chg.txt
share/doc/ftp4all/f4adp/usr_get.txt
share/doc/ftp4all/f4adp/usr_set.txt
share/doc/ftp4all/f4adp/usr_kic.txt
share/doc/ftp4all/f4adp/usr_kil.txt
share/doc/ftp4all/f4adp/usr_ls.txt
share/doc/ftp4all/f4adp/usr_lsr.txt
share/doc/ftp4all/f4adp/usr_pwd.txt
share/doc/ftp4all/f4adp/usr_rm.txt
share/doc/ftp4all/f4adp/usr_rst.txt
share/doc/ftp4all/f4adp/var_exp.txt
share/doc/ftp4all/f4adp/ver_f4p.txt
share/doc/ftp4all/f4adp/wa_add.txt
share/doc/ftp4all/f4adp/adi.txt
share/doc/ftp4all/frames/faq_user.html
share/doc/ftp4all/frames/configs.html
share/doc/ftp4all/frames/copyright.html
share/doc/ftp4all/frames/debug.html
share/doc/ftp4all/frames/faq_v3.html
share/doc/ftp4all/frames/history.html
share/doc/ftp4all/frames/index_v3.html
share/doc/ftp4all/frames/install.html
share/doc/ftp4all/frames/new.html
share/doc/ftp4all/frames/registry.html
share/doc/ftp4all/frames/sitecmds.html
share/doc/ftp4all/frames/upgrade.html
share/doc/ftp4all/frames/index.html
share/doc/ftp4all/frames/menu.html
share/doc/ftp4all/frames/credits.html
share/doc/ftp4all/frames/download.html
share/doc/ftp4all/frames/f4agui.html
share/doc/ftp4all/frames/messages.html
share/doc/ftp4all/frames/faq_config.html
share/doc/ftp4all/frames/faq_f4agui.html
share/doc/ftp4all/frames/faq_files.html
share/doc/ftp4all/frames/faq_index.html
share/doc/ftp4all/frames/faq_start.html
share/doc/ftp4all/frames/faq_system.html
share/doc/ftp4all/frames/concepts.html
share/doc/ftp4all/noframes/faq_user.html
share/doc/ftp4all/noframes/configs.html
share/doc/ftp4all/noframes/copyright.html
share/doc/ftp4all/noframes/debug.html
share/doc/ftp4all/noframes/faq_v3.html
share/doc/ftp4all/noframes/history.html
share/doc/ftp4all/noframes/index_v3.html
share/doc/ftp4all/noframes/install.html
share/doc/ftp4all/noframes/new.html
share/doc/ftp4all/noframes/registry.html
share/doc/ftp4all/noframes/sitecmds.html
share/doc/ftp4all/noframes/upgrade.html
share/doc/ftp4all/noframes/index.html
share/doc/ftp4all/noframes/menu.html
share/doc/ftp4all/noframes/credits.html
share/doc/ftp4all/noframes/download.html
share/doc/ftp4all/noframes/f4agui.html
share/doc/ftp4all/noframes/messages.html
share/doc/ftp4all/noframes/faq_config.html
share/doc/ftp4all/noframes/faq_f4agui.html
share/doc/ftp4all/noframes/faq_files.html
share/doc/ftp4all/noframes/faq_index.html
share/doc/ftp4all/noframes/faq_start.html
share/doc/ftp4all/noframes/faq_system.html
share/doc/ftp4all/noframes/concepts.html
@dirrm etc/ftp4all/msg/grp
@dirrm etc/ftp4all/msg/usr
@dirrm etc/ftp4all/msg
@dirrm etc/ftp4all
@dirrm share/doc/ftp4all/f4adp
@dirrm share/doc/ftp4all/frames
@dirrm share/doc/ftp4all/noframes
@dirrm share/doc/ftp4all
