# net/pipsecd/files/pipsecd.in
#!/bin/sh
#
#$FreeBSD$
#
# PROVIDE: pipsecd
# REQUIRE: root mountcritlocal
# BEFORE:  DAEMON
# KEYWORD: nojail

. /etc/rc.subr

name="pipsecd"
rcvar=pipsecd_enable
command="%%PREFIX%%/sbin/${name}"
command_args="${pipsecd_flags}"
start_cmd="pipsecd_start"

: ${pipsecd_enable="NO"}

pipsecd_start() {
	if ! checkyesno pipsecd_enable ; then
		return 0
	fi
	echo "Starting ${name}."
	${command} ${pipsecd_flags} &
}

load_rc_config $name
run_rc_command "$1"
