# sysutils/nwclient602/pkg-plist
etc/rc.d/networker.sh
nsr/bin/backint
nsr/bin/lsm_backint
nsr/bin/mminfo
nsr/bin/mmlocate
nsr/bin/mmpool
nsr/bin/mmrecov
nsr/bin/nsradmin
nsr/bin/nsralist
nsr/bin/nsrarchive
nsr/bin/nsrclone
nsr/bin/nsrexec
nsr/bin/nsrexecd
nsr/bin/nsrfile
nsr/bin/nsrifmx
nsr/bin/nsrmm
nsr/bin/nsrmon
nsr/bin/nsrports
nsr/bin/nsrretrieve
nsr/bin/nsrssc
nsr/bin/nsrstage
nsr/bin/nsrwatch
nsr/bin/preclntsave
nsr/bin/pstclntsave
nsr/bin/recover
nsr/bin/sapclone
nsr/bin/save
nsr/bin/savefs
nsr/bin/savegrp
nsr/bin/savepnpc
nsr/bin/uasm
nsr/gls/cm/registry
nsr/gls/lc/os/portable/C
man/cat8/ansrd.8
man/cat8/changers.8
man/cat8/dasadmin.8
man/cat8/EMASS_silo.8
man/cat8/emcdiscover.8
man/cat8/erase.8
man/cat8/hpflip.8
man/cat8/IBM_silo.8
man/cat8/ielem.8
man/cat8/inquire.8
man/cat8/jbconfig.8
man/cat8/jbexercise.8
man/cat8/jbinfo.8
man/cat8/ldunld.8
man/cat8/lgtolic.8
man/cat8/lgtolmd.8
man/cat8/libscsi.8
man/cat8/libsji.8
man/cat8/libstlemass.8
man/cat8/libstlibm.8
man/cat8/libstlstk.8
man/cat8/lrescan.8
man/cat8/lreset.8
man/cat8/lusbinfo.8
man/cat8/lusdebug.8
man/cat8/lusmode.8
man/cat8/mini_el.8
man/cat8/mminfo.8
man/cat8/mmlocate.8
man/cat8/mmpool.8
man/cat8/mmrecov.8
man/cat8/msense.8
man/cat8/networker.8
man/cat8/nmi_config.8
man/cat8/nsr.8
man/cat8/nsr_crash.8
man/cat8/nsr_ize.8
man/cat8/nsr_shutdown.8
man/cat8/nsradmin.8
man/cat8/nsralist.8
man/cat8/nsrarchive.8
man/cat8/nsrcap.8
man/cat8/nsrcat.8
man/cat8/nsrck.8
man/cat8/nsrclone.8
man/cat8/nsrd.8
man/cat8/nsrdb2.8
man/cat8/nsrdbmi.8
man/cat8/nsrdocrc.8
man/cat8/nsrexec.8
man/cat8/nsrexecd.8
man/cat8/nsrhsmck.8
man/cat8/nsrhsmls.8
man/cat8/nsrib.8
man/cat8/nsriba.8
man/cat8/nsrim.8
man/cat8/nsrindexasm.8
man/cat8/nsrindexd.8
man/cat8/nsrinfo.8
man/cat8/nsrjb.8
man/cat8/nsrlic.8
man/cat8/nsrls.8
man/cat8/nsrmig.8
man/cat8/nsrmm.8
man/cat8/nsrmmd.8
man/cat8/nsrmmdbasm.8
man/cat8/nsrmmdbd.8
man/cat8/nsrmon.8
man/cat8/nsrndmp_recover.8
man/cat8/nsrndmp_save.8
man/cat8/nsrnmo.8
man/cat8/nsrnotesrc.8
man/cat8/nsrnotesv.8
man/cat8/nsrpmig.8
man/cat8/nsrports.8
man/cat8/nsrretrieve.8
man/cat8/nsrssc.8
man/cat8/nsrstage.8
man/cat8/nsrsyb.8
man/cat8/nsrsybcc.8
man/cat8/nsrsybrc.8
man/cat8/nsrsybsv.8
man/cat8/nsrtrap.8
man/cat8/nsrwatch.8
man/cat8/nwadmin.8
man/cat8/nwarchive.8
man/cat8/nwbackup.8
man/cat8/nwrecover.8
man/cat8/nwretrieve.8
man/cat8/oraemcasm.8
man/cat8/oraemcmap.8
man/cat8/pmode.8
man/cat8/preclntsave.8
man/cat8/pstclntsave.8
man/cat8/rap.8
man/cat8/rapcheck.8
man/cat8/rapd.8
man/cat8/rapxfer.8
man/cat8/recover.8
man/cat8/relem.8
man/cat8/sapclone.8
man/cat8/save.8
man/cat8/savefs.8
man/cat8/savegrp.8
man/cat8/savepnpc.8
man/cat8/scanner.8
man/cat8/sjidopen.8
man/cat8/sjiielm.8
man/cat8/sjiinq.8
man/cat8/sjirdp.8
man/cat8/sjirdtag.8
man/cat8/sjirelem.8
man/cat8/sjirjc.8
man/cat8/ssi.8
man/cat8/stk_eject.8
man/cat8/STK_silo.8
man/cat8/stli.8
man/cat8/tapeexercise.8
man/cat8/tur.8
man/cat8/uasm.8
man/cat8/writebuf.8
man/cat5/mm_data.5
man/cat5/nsr.5
man/cat5/nsr_archive_request.5
man/cat5/nsr_client.5
man/cat5/nsr_data.5
man/cat5/nsr_device.5
man/cat5/nsr_directive.5
man/cat5/nsr_group.5
man/cat5/nsr_jukebox.5
man/cat5/nsr_label.5
man/cat5/nsr_layout.5
man/cat5/nsr_license.5
man/cat5/nsr_migration.5
man/cat5/nsr_notification.5
man/cat5/nsr_policy.5
man/cat5/nsr_pool.5
man/cat5/nsr_regexp.5
man/cat5/nsr_resource.5
man/cat5/nsr_schedule.5
man/cat5/nsr_service.5
man/cat5/nsr_stage.5
man/cat5/nsr_storage_node.5
man/cat5/resource.5
man/cat1/rapgen.1
man/cat3/nsr_getdate.3
@dirrm nsr/bin
@dirrm nsr/gls/cm
@dirrm nsr/gls/lc/os/portable
@dirrm nsr/gls/lc/os
@dirrm nsr/gls/lc
@dirrm nsr/gls
@dirrm nsr
