# print/cups-base/pkg-plist
@exec /bin/echo "===>" Stopping cupsd ...
@exec /usr/bin/killall cupsd || true
@exec /bin/rm -fr %D/etc/cups/certs/* 2>/dev/null || true
bin/cups-config
@unexec if cmp -s %D/etc/cups/classes.conf.N %D/etc/cups/classes.conf; then rm -f %D/etc/cups/classes.conf; fi
etc/cups/classes.conf.N
@exec if test ! -f %D/etc/cups/classes.conf; then cp -p %D/etc/cups/classes.conf.N %D/etc/cups/classes.conf; fi
@unexec if cmp -s %D/etc/cups/client.conf.N %D/etc/cups/client.conf; then rm -f %D/etc/cups/client.conf; fi
etc/cups/client.conf.N
@exec if test ! -f %D/etc/cups/client.conf; then cp -p %D/etc/cups/client.conf.N %D/etc/cups/client.conf; fi
@unexec if cmp -s %D/etc/cups/cupsd.conf.N %D/etc/cups/cupsd.conf; then rm -f %D/etc/cups/cupsd.conf; fi
etc/cups/cupsd.conf.N
@exec if test ! -f %D/etc/cups/cupsd.conf; then cp -p %D/etc/cups/cupsd.conf.N %D/etc/cups/cupsd.conf; fi
etc/cups/mime.convs
etc/cups/mime.types
@unexec if cmp -s %D/etc/cups/printers.conf.N %D/etc/cups/printers.conf; then rm -f %D/etc/cups/printers.conf; fi
etc/cups/printers.conf.N
@exec if test ! -f %D/etc/cups/printers.conf; then cp -p %D/etc/cups/printers.conf.N %D/etc/cups/printers.conf; fi
etc/rc.d/cups.sh.sample
include/cups/cups.h
include/cups/http.h
include/cups/image.h
include/cups/ipp.h
include/cups/language.h
include/cups/md5.h
include/cups/ppd.h
include/cups/raster.h
lib/libcups.a
lib/libcups.so
lib/libcups.so.2
lib/libcupsimage.a
lib/libcupsimage.so
lib/libcupsimage.so.2
libexec/cups/backend/http
libexec/cups/backend/ipp
libexec/cups/backend/lpd
libexec/cups/backend/parallel
libexec/cups/backend/scsi
libexec/cups/backend/serial
libexec/cups/backend/socket
libexec/cups/backend/usb
libexec/cups/cgi-bin/admin.cgi
libexec/cups/cgi-bin/classes.cgi
libexec/cups/cgi-bin/jobs.cgi
libexec/cups/cgi-bin/printers.cgi
libexec/cups/daemon/cups-lpd
libexec/cups/daemon/cups-polld
libexec/cups/filter/gziptoany
libexec/cups/filter/hpgltops
libexec/cups/filter/imagetops
libexec/cups/filter/imagetoraster
libexec/cups/filter/pdftops
libexec/cups/filter/pstops
libexec/cups/filter/rastertodymo
libexec/cups/filter/rastertoepson
libexec/cups/filter/rastertohp
libexec/cups/filter/texttops
sbin/cupsd
share/cups/banners/classified
share/cups/banners/confidential
share/cups/banners/secret
share/cups/banners/standard
share/cups/banners/topsecret
share/cups/banners/unclassified
share/cups/charsets/iso-8859-1
share/cups/charsets/iso-8859-10
share/cups/charsets/iso-8859-13
share/cups/charsets/iso-8859-14
share/cups/charsets/iso-8859-15
share/cups/charsets/iso-8859-2
share/cups/charsets/iso-8859-3
share/cups/charsets/iso-8859-4
share/cups/charsets/iso-8859-5
share/cups/charsets/iso-8859-6
share/cups/charsets/iso-8859-7
share/cups/charsets/iso-8859-8
share/cups/charsets/iso-8859-9
share/cups/charsets/koi8-r
share/cups/charsets/koi8-u
share/cups/charsets/utf-8
share/cups/charsets/windows-1250
share/cups/charsets/windows-1251
share/cups/charsets/windows-1252
share/cups/charsets/windows-1253
share/cups/charsets/windows-1254
share/cups/charsets/windows-1255
share/cups/charsets/windows-1256
share/cups/charsets/windows-1257
share/cups/charsets/windows-1258
share/cups/charsets/windows-874
share/cups/data/HPGLprolog
share/cups/data/psglyphs
share/cups/data/testprint.ps
share/cups/fonts/Courier
share/cups/fonts/Courier-Bold
share/cups/fonts/Courier-BoldOblique
share/cups/fonts/Courier-Oblique
share/cups/fonts/Symbol
share/cups/model/deskjet.ppd
share/cups/model/deskjet2.ppd
share/cups/model/dymo.ppd
share/cups/model/epson24.ppd
share/cups/model/epson9.ppd
share/cups/model/laserjet.ppd
share/cups/model/okidat24.ppd
share/cups/model/okidata9.ppd
share/cups/model/stcolor.ppd
share/cups/model/stcolor2.ppd
share/cups/model/stphoto.ppd
share/cups/model/stphoto2.ppd
share/cups/templates/add-class.tmpl
share/cups/templates/add-printer.tmpl
share/cups/templates/admin-op.tmpl
share/cups/templates/admin.tmpl
share/cups/templates/choose-device.tmpl
share/cups/templates/choose-make.tmpl
share/cups/templates/choose-members.tmpl
share/cups/templates/choose-model.tmpl
share/cups/templates/choose-serial.tmpl
share/cups/templates/choose-uri.tmpl
share/cups/templates/class-added.tmpl
share/cups/templates/class-confirm.tmpl
share/cups/templates/class-deleted.tmpl
share/cups/templates/class-modified.tmpl
share/cups/templates/classes.tmpl
share/cups/templates/config-printer.tmpl
share/cups/templates/config-printer2.tmpl
share/cups/templates/de/add-class.tmpl
share/cups/templates/de/add-printer.tmpl
share/cups/templates/de/admin-op.tmpl
share/cups/templates/de/admin.tmpl
share/cups/templates/de/choose-device.tmpl
share/cups/templates/de/choose-make.tmpl
share/cups/templates/de/choose-members.tmpl
share/cups/templates/de/choose-model.tmpl
share/cups/templates/de/choose-serial.tmpl
share/cups/templates/de/choose-uri.tmpl
share/cups/templates/de/class-added.tmpl
share/cups/templates/de/class-confirm.tmpl
share/cups/templates/de/class-deleted.tmpl
share/cups/templates/de/class-modified.tmpl
share/cups/templates/de/classes.tmpl
share/cups/templates/de/config-printer.tmpl
share/cups/templates/de/config-printer2.tmpl
share/cups/templates/de/error.tmpl
share/cups/templates/de/header.tmpl
share/cups/templates/de/job-cancel.tmpl
share/cups/templates/de/job-hold.tmpl
share/cups/templates/de/job-op.tmpl
share/cups/templates/de/job-release.tmpl
share/cups/templates/de/job-restart.tmpl
share/cups/templates/de/jobs.tmpl
share/cups/templates/de/modify-class.tmpl
share/cups/templates/de/modify-printer.tmpl
share/cups/templates/de/option-boolean.tmpl
share/cups/templates/de/option-conflict.tmpl
share/cups/templates/de/option-header.tmpl
share/cups/templates/de/option-pickmany.tmpl
share/cups/templates/de/option-pickone.tmpl
share/cups/templates/de/option-trailer.tmpl
share/cups/templates/de/printer-accept.tmpl
share/cups/templates/de/printer-added.tmpl
share/cups/templates/de/printer-configured.tmpl
share/cups/templates/de/printer-confirm.tmpl
share/cups/templates/de/printer-deleted.tmpl
share/cups/templates/de/printer-modified.tmpl
share/cups/templates/de/printer-purge.tmpl
share/cups/templates/de/printer-reject.tmpl
share/cups/templates/de/printer-start.tmpl
share/cups/templates/de/printer-stop.tmpl
share/cups/templates/de/printers.tmpl
share/cups/templates/de/test-page.tmpl
share/cups/templates/de/trailer.tmpl
share/cups/templates/error.tmpl
share/cups/templates/fr/add-class.tmpl
share/cups/templates/fr/add-printer.tmpl
share/cups/templates/fr/admin-op.tmpl
share/cups/templates/fr/admin.tmpl
share/cups/templates/fr/choose-device.tmpl
share/cups/templates/fr/choose-make.tmpl
share/cups/templates/fr/choose-members.tmpl
share/cups/templates/fr/choose-model.tmpl
share/cups/templates/fr/choose-serial.tmpl
share/cups/templates/fr/choose-uri.tmpl
share/cups/templates/fr/class-added.tmpl
share/cups/templates/fr/class-confirm.tmpl
share/cups/templates/fr/class-deleted.tmpl
share/cups/templates/fr/class-modified.tmpl
share/cups/templates/fr/classes.tmpl
share/cups/templates/fr/config-printer.tmpl
share/cups/templates/fr/config-printer2.tmpl
share/cups/templates/fr/error.tmpl
share/cups/templates/fr/header.tmpl
share/cups/templates/fr/job-cancel.tmpl
share/cups/templates/fr/job-hold.tmpl
share/cups/templates/fr/job-op.tmpl
share/cups/templates/fr/job-release.tmpl
share/cups/templates/fr/job-restart.tmpl
share/cups/templates/fr/jobs.tmpl
share/cups/templates/fr/modify-class.tmpl
share/cups/templates/fr/modify-printer.tmpl
share/cups/templates/fr/option-boolean.tmpl
share/cups/templates/fr/option-conflict.tmpl
share/cups/templates/fr/option-header.tmpl
share/cups/templates/fr/option-pickmany.tmpl
share/cups/templates/fr/option-pickone.tmpl
share/cups/templates/fr/option-trailer.tmpl
share/cups/templates/fr/printer-accept.tmpl
share/cups/templates/fr/printer-added.tmpl
share/cups/templates/fr/printer-configured.tmpl
share/cups/templates/fr/printer-confirm.tmpl
share/cups/templates/fr/printer-deleted.tmpl
share/cups/templates/fr/printer-modified.tmpl
share/cups/templates/fr/printer-purge.tmpl
share/cups/templates/fr/printer-reject.tmpl
share/cups/templates/fr/printer-start.tmpl
share/cups/templates/fr/printer-stop.tmpl
share/cups/templates/fr/printers.tmpl
share/cups/templates/fr/test-page.tmpl
share/cups/templates/fr/trailer.tmpl
share/cups/templates/header.tmpl
share/cups/templates/job-cancel.tmpl
share/cups/templates/job-hold.tmpl
share/cups/templates/job-op.tmpl
share/cups/templates/job-release.tmpl
share/cups/templates/job-restart.tmpl
share/cups/templates/jobs.tmpl
share/cups/templates/modify-class.tmpl
share/cups/templates/modify-printer.tmpl
share/cups/templates/option-boolean.tmpl
share/cups/templates/option-conflict.tmpl
share/cups/templates/option-header.tmpl
share/cups/templates/option-pickmany.tmpl
share/cups/templates/option-pickone.tmpl
share/cups/templates/option-trailer.tmpl
share/cups/templates/printer-accept.tmpl
share/cups/templates/printer-added.tmpl
share/cups/templates/printer-configured.tmpl
share/cups/templates/printer-confirm.tmpl
share/cups/templates/printer-deleted.tmpl
share/cups/templates/printer-modified.tmpl
share/cups/templates/printer-purge.tmpl
share/cups/templates/printer-reject.tmpl
share/cups/templates/printer-start.tmpl
share/cups/templates/printer-stop.tmpl
share/cups/templates/printers.tmpl
share/cups/templates/test-page.tmpl
share/cups/templates/trailer.tmpl
share/doc/cups/cmp.html
share/doc/cups/cmp.pdf
share/doc/cups/cups.css
share/doc/cups/cupsdoc.css
share/doc/cups/de/documentation.html
share/doc/cups/de/images/accept-jobs.gif
share/doc/cups/de/images/add-class.gif
share/doc/cups/de/images/add-printer.gif
share/doc/cups/de/images/cancel-job.gif
share/doc/cups/de/images/config-printer.gif
share/doc/cups/de/images/continue.gif
share/doc/cups/de/images/delete-class.gif
share/doc/cups/de/images/delete-printer.gif
share/doc/cups/de/images/hold-job.gif
share/doc/cups/de/images/manage-classes.gif
share/doc/cups/de/images/manage-jobs.gif
share/doc/cups/de/images/manage-printers.gif
share/doc/cups/de/images/modify-class.gif
share/doc/cups/de/images/modify-printer.gif
share/doc/cups/de/images/navbar.gif
share/doc/cups/de/images/print-test-page.gif
share/doc/cups/de/images/reject-jobs.gif
share/doc/cups/de/images/release-job.gif
share/doc/cups/de/images/restart-job.gif
share/doc/cups/de/images/show-active.gif
share/doc/cups/de/images/show-completed.gif
share/doc/cups/de/images/start-class.gif
share/doc/cups/de/images/start-printer.gif
share/doc/cups/de/images/stop-class.gif
share/doc/cups/de/images/stop-printer.gif
share/doc/cups/de/index.html
share/doc/cups/documentation.html
share/doc/cups/fr/cups.css
share/doc/cups/fr/cupsdoc.css
share/doc/cups/fr/documentation.html
share/doc/cups/fr/images/accept-jobs.gif
share/doc/cups/fr/images/add-class.gif
share/doc/cups/fr/images/add-printer.gif
share/doc/cups/fr/images/cancel-job.gif
share/doc/cups/fr/images/config-printer.gif
share/doc/cups/fr/images/continue.gif
share/doc/cups/fr/images/delete-class.gif
share/doc/cups/fr/images/delete-printer.gif
share/doc/cups/fr/images/hold-job.gif
share/doc/cups/fr/images/manage-classes.gif
share/doc/cups/fr/images/manage-jobs.gif
share/doc/cups/fr/images/manage-printers.gif
share/doc/cups/fr/images/modify-class.gif
share/doc/cups/fr/images/modify-printer.gif
share/doc/cups/fr/images/navbar.gif
share/doc/cups/fr/images/print-test-page.gif
share/doc/cups/fr/images/reject-jobs.gif
share/doc/cups/fr/images/release-job.gif
share/doc/cups/fr/images/restart-job.gif
share/doc/cups/fr/images/show-active.gif
share/doc/cups/fr/images/show-completed.gif
share/doc/cups/fr/images/start-class.gif
share/doc/cups/fr/images/start-printer.gif
share/doc/cups/fr/images/stop-class.gif
share/doc/cups/fr/images/stop-printer.gif
share/doc/cups/fr/index.html
share/doc/cups/fr/overview.html
share/doc/cups/fr/overview.pdf
share/doc/cups/fr/sam.html
share/doc/cups/fr/sam.pdf
share/doc/cups/fr/sum.html
share/doc/cups/fr/sum.pdf
share/doc/cups/idd.html
share/doc/cups/idd.pdf
share/doc/cups/images/accept-jobs.gif
share/doc/cups/images/add-class.gif
share/doc/cups/images/add-printer.gif
share/doc/cups/images/cancel-job.gif
share/doc/cups/images/classes.gif
share/doc/cups/images/config-printer.gif
share/doc/cups/images/continue.gif
share/doc/cups/images/cups-block-diagram.gif
share/doc/cups/images/cups-large.gif
share/doc/cups/images/cups-medium.gif
share/doc/cups/images/cups-small.gif
share/doc/cups/images/delete-class.gif
share/doc/cups/images/delete-printer.gif
share/doc/cups/images/hold-job.gif
share/doc/cups/images/left.gif
share/doc/cups/images/logo.gif
share/doc/cups/images/manage-classes.gif
share/doc/cups/images/manage-jobs.gif
share/doc/cups/images/manage-printers.gif
share/doc/cups/images/modify-class.gif
share/doc/cups/images/modify-printer.gif
share/doc/cups/images/navbar.gif
share/doc/cups/images/print-test-page.gif
share/doc/cups/images/printer-idle.gif
share/doc/cups/images/printer-processing.gif
share/doc/cups/images/printer-stopped.gif
share/doc/cups/images/reject-jobs.gif
share/doc/cups/images/release-job.gif
share/doc/cups/images/restart-job.gif
share/doc/cups/images/right.gif
share/doc/cups/images/show-active.gif
share/doc/cups/images/show-completed.gif
share/doc/cups/images/start-class.gif
share/doc/cups/images/start-printer.gif
share/doc/cups/images/stop-class.gif
share/doc/cups/images/stop-printer.gif
share/doc/cups/index.html
share/doc/cups/ipp.html
share/doc/cups/ipp.pdf
share/doc/cups/overview.html
share/doc/cups/overview.pdf
share/doc/cups/robots.txt
share/doc/cups/sam.html
share/doc/cups/sam.pdf
share/doc/cups/sdd.html
share/doc/cups/sdd.pdf
share/doc/cups/spm.html
share/doc/cups/spm.pdf
share/doc/cups/sps.html
share/doc/cups/sps.pdf
share/doc/cups/ssr.html
share/doc/cups/ssr.pdf
share/doc/cups/stp.html
share/doc/cups/stp.pdf
share/doc/cups/sum.html
share/doc/cups/sum.pdf
share/doc/cups/svd.html
share/doc/cups/svd.pdf
share/doc/cups/translation.html
share/doc/cups/translation.pdf
share/locale/C/cups_C
share/locale/be/cups_be
share/locale/cs/cups_cs
share/locale/de/cups_de
share/locale/en/cups_en
share/locale/en_US/cups_en_US
share/locale/es/cups_es
share/locale/fr/cups_fr
share/locale/he/cups_he
share/locale/it/cups_it
share/locale/ru_RU/cups_ru_RU
share/locale/sv/cups_sv
share/locale/uk/cups_uk
share/locale/uk_UA/cups_uk_UA
share/locale/zh_CN/cups_zh_CN
@unexec rmdir %D/share/locale/zh_CN 2>/dev/null || true
@unexec rmdir %D/share/locale/uk_UA 2>/dev/null || true
@unexec rmdir %D/share/locale/uk 2>/dev/null || true
@unexec rmdir %D/share/locale/sv 2>/dev/null || true
@unexec rmdir %D/share/locale/ru_RU 2>/dev/null || true
@unexec rmdir %D/share/locale/it 2>/dev/null || true
@unexec rmdir %D/share/locale/he 2>/dev/null || true
@unexec rmdir %D/share/locale/fr 2>/dev/null || true
@unexec rmdir %D/share/locale/es 2>/dev/null || true
@unexec rmdir %D/share/locale/en 2>/dev/null || true
@unexec rmdir %D/share/locale/en_US 2>/dev/null || true
@unexec rmdir %D/share/locale/de 2>/dev/null || true
@unexec rmdir %D/share/locale/cs 2>/dev/null || true
@unexec rmdir %D/share/locale/be 2>/dev/null || true
@unexec rmdir %D/share/locale/C 2>/dev/null || true
@unexec rmdir %D/share/locale 2>/dev/null || true
@dirrm share/doc/cups/images
@dirrm share/doc/cups/fr/images
@dirrm share/doc/cups/fr
@dirrm share/doc/cups/de/images
@dirrm share/doc/cups/de
@dirrm share/doc/cups
@dirrm share/cups/templates/fr
@dirrm share/cups/templates/de
@dirrm share/cups/templates
@dirrm share/cups/model
@dirrm share/cups/fonts
@dirrm share/cups/data
@dirrm share/cups/charsets
@dirrm share/cups/banners
@dirrm share/cups
@dirrm libexec/cups/filter
@dirrm libexec/cups/daemon
@dirrm libexec/cups/cgi-bin
@dirrm libexec/cups/backend
@dirrm libexec/cups
@dirrm include/cups
@dirrm etc/cups/ppd
@dirrm etc/cups/interfaces
@dirrm etc/cups/certs
@exec mkdir -p /var/spool/cups/tmp || true
@exec chown -R root:daemon /var/spool/cups
@exec mkdir %D/etc/cups/certs || true
@exec mkdir %D/etc/cups/ppd || true
@exec mkdir %D/etc/cups/interfaces || true
@unexec rm -rf /var/spool/cups
@unexec rmdir %D/etc/cups 2>/dev/null || true
