# print/cups-base/pkg-plist
bin/cups-config
etc/cups/certs/0
etc/cups/classes.conf.N
etc/cups/client.conf.N
etc/cups/cupsd.conf.N
etc/cups/mime.convs
etc/cups/mime.types
etc/cups/printers.conf.N
etc/rc.d/cups.sh.sample
include/cups/cups.h
include/cups/http.h
include/cups/image.h
include/cups/ipp.h
include/cups/language.h
include/cups/md5.h
include/cups/ppd.h
include/cups/raster.h
lib/libcups.a
lib/libcups.so
lib/libcups.so.2
lib/libcupsimage.a
lib/libcupsimage.so
lib/libcupsimage.so.2
libexec/cups/backend/http
libexec/cups/backend/ipp
libexec/cups/backend/lpd
libexec/cups/backend/parallel
libexec/cups/backend/scsi
libexec/cups/backend/serial
libexec/cups/backend/socket
libexec/cups/backend/usb
libexec/cups/cgi-bin/admin.cgi
libexec/cups/cgi-bin/classes.cgi
libexec/cups/cgi-bin/jobs.cgi
libexec/cups/cgi-bin/printers.cgi
libexec/cups/daemon/cups-lpd
libexec/cups/daemon/cups-polld
libexec/cups/filter/hpgltops
libexec/cups/filter/imagetops
libexec/cups/filter/imagetoraster
libexec/cups/filter/pdftops
libexec/cups/filter/pstops
libexec/cups/filter/rastertodymo
libexec/cups/filter/rastertoepson
libexec/cups/filter/rastertohp
libexec/cups/filter/texttops
sbin/cupsd
share/cups/banners/classified
share/cups/banners/confidential
share/cups/banners/secret
share/cups/banners/standard
share/cups/banners/topsecret
share/cups/banners/unclassified
share/cups/charsets/iso-8859-1
share/cups/charsets/iso-8859-10
share/cups/charsets/iso-8859-13
share/cups/charsets/iso-8859-14
share/cups/charsets/iso-8859-15
share/cups/charsets/iso-8859-2
share/cups/charsets/iso-8859-3
share/cups/charsets/iso-8859-4
share/cups/charsets/iso-8859-5
share/cups/charsets/iso-8859-6
share/cups/charsets/iso-8859-7
share/cups/charsets/iso-8859-8
share/cups/charsets/iso-8859-9
share/cups/charsets/koi8-r
share/cups/charsets/koi8-u
share/cups/charsets/utf-8
share/cups/charsets/windows-1250
share/cups/charsets/windows-1251
share/cups/charsets/windows-1252
share/cups/charsets/windows-1253
share/cups/charsets/windows-1254
share/cups/charsets/windows-1255
share/cups/charsets/windows-1256
share/cups/charsets/windows-1257
share/cups/charsets/windows-1258
share/cups/charsets/windows-874
share/cups/data/HPGLprolog
share/cups/data/psglyphs
share/cups/data/testprint.ps
share/cups/fonts/Courier
share/cups/fonts/Courier-Bold
share/cups/fonts/Courier-BoldOblique
share/cups/fonts/Courier-Oblique
share/cups/fonts/Symbol
share/cups/model/deskjet.ppd
share/cups/model/deskjet2.ppd
share/cups/model/dymo.ppd
share/cups/model/epson24.ppd
share/cups/model/epson9.ppd
share/cups/model/laserjet.ppd
share/cups/model/okidat24.ppd
share/cups/model/okidata9.ppd
share/cups/model/stcolor.ppd
share/cups/model/stcolor2.ppd
share/cups/model/stphoto.ppd
share/cups/model/stphoto2.ppd
share/cups/templates/add-class.tmpl
share/cups/templates/add-printer.tmpl
share/cups/templates/admin-op.tmpl
share/cups/templates/admin.tmpl
share/cups/templates/choose-device.tmpl
share/cups/templates/choose-make.tmpl
share/cups/templates/choose-members.tmpl
share/cups/templates/choose-model.tmpl
share/cups/templates/choose-serial.tmpl
share/cups/templates/choose-uri.tmpl
share/cups/templates/class-added.tmpl
share/cups/templates/class-confirm.tmpl
share/cups/templates/class-deleted.tmpl
share/cups/templates/class-modified.tmpl
share/cups/templates/classes.tmpl
share/cups/templates/config-printer.tmpl
share/cups/templates/config-printer2.tmpl
share/cups/templates/error.tmpl
share/cups/templates/fr/add-class.tmpl
share/cups/templates/fr/add-printer.tmpl
share/cups/templates/fr/admin-op.tmpl
share/cups/templates/fr/admin.tmpl
share/cups/templates/fr/choose-device.tmpl
share/cups/templates/fr/choose-make.tmpl
share/cups/templates/fr/choose-members.tmpl
share/cups/templates/fr/choose-model.tmpl
share/cups/templates/fr/choose-serial.tmpl
share/cups/templates/fr/choose-uri.tmpl
share/cups/templates/fr/class-added.tmpl
share/cups/templates/fr/class-confirm.tmpl
share/cups/templates/fr/class-deleted.tmpl
share/cups/templates/fr/class-modified.tmpl
share/cups/templates/fr/classes.tmpl
share/cups/templates/fr/config-printer.tmpl
share/cups/templates/fr/config-printer2.tmpl
share/cups/templates/fr/error.tmpl
share/cups/templates/fr/header.tmpl
share/cups/templates/fr/job-cancel.tmpl
share/cups/templates/fr/job-hold.tmpl
share/cups/templates/fr/job-op.tmpl
share/cups/templates/fr/job-release.tmpl
share/cups/templates/fr/job-restart.tmpl
share/cups/templates/fr/jobs.tmpl
share/cups/templates/fr/modify-class.tmpl
share/cups/templates/fr/modify-printer.tmpl
share/cups/templates/fr/option-boolean.tmpl
share/cups/templates/fr/option-header.tmpl
share/cups/templates/fr/option-pickmany.tmpl
share/cups/templates/fr/option-pickone.tmpl
share/cups/templates/fr/option-trailer.tmpl
share/cups/templates/fr/printer-accept.tmpl
share/cups/templates/fr/printer-added.tmpl
share/cups/templates/fr/printer-configured.tmpl
share/cups/templates/fr/printer-confirm.tmpl
share/cups/templates/fr/printer-deleted.tmpl
share/cups/templates/fr/printer-modified.tmpl
share/cups/templates/fr/printer-purge.tmpl
share/cups/templates/fr/printer-reject.tmpl
share/cups/templates/fr/printer-start.tmpl
share/cups/templates/fr/printer-stop.tmpl
share/cups/templates/fr/printers.tmpl
share/cups/templates/fr/test-page.tmpl
share/cups/templates/fr/trailer.tmpl
share/cups/templates/header.tmpl
share/cups/templates/job-cancel.tmpl
share/cups/templates/job-hold.tmpl
share/cups/templates/job-op.tmpl
share/cups/templates/job-release.tmpl
share/cups/templates/job-restart.tmpl
share/cups/templates/jobs.tmpl
share/cups/templates/modify-class.tmpl
share/cups/templates/modify-printer.tmpl
share/cups/templates/option-boolean.tmpl
share/cups/templates/option-conflict.tmpl
share/cups/templates/option-header.tmpl
share/cups/templates/option-pickmany.tmpl
share/cups/templates/option-pickone.tmpl
share/cups/templates/option-trailer.tmpl
share/cups/templates/printer-accept.tmpl
share/cups/templates/printer-added.tmpl
share/cups/templates/printer-configured.tmpl
share/cups/templates/printer-confirm.tmpl
share/cups/templates/printer-deleted.tmpl
share/cups/templates/printer-modified.tmpl
share/cups/templates/printer-purge.tmpl
share/cups/templates/printer-reject.tmpl
share/cups/templates/printer-start.tmpl
share/cups/templates/printer-stop.tmpl
share/cups/templates/printers.tmpl
share/cups/templates/test-page.tmpl
share/cups/templates/trailer.tmpl
share/doc/cups-1.1.18.0/cmp.html
share/doc/cups-1.1.18.0/cmp.pdf
share/doc/cups-1.1.18.0/cups.css
share/doc/cups-1.1.18.0/cupsdoc.css
share/doc/cups-1.1.18.0/documentation.html
share/doc/cups-1.1.18.0/fr/cups.css
share/doc/cups-1.1.18.0/fr/cupsdoc.css
share/doc/cups-1.1.18.0/fr/documentation.html
share/doc/cups-1.1.18.0/fr/images/accept-jobs.gif
share/doc/cups-1.1.18.0/fr/images/add-class.gif
share/doc/cups-1.1.18.0/fr/images/add-printer.gif
share/doc/cups-1.1.18.0/fr/images/cancel-job.gif
share/doc/cups-1.1.18.0/fr/images/cancel-jobs.gif
share/doc/cups-1.1.18.0/fr/images/cancel.gif
share/doc/cups-1.1.18.0/fr/images/config-printer.gif
share/doc/cups-1.1.18.0/fr/images/continue.gif
share/doc/cups-1.1.18.0/fr/images/delete-class.gif
share/doc/cups-1.1.18.0/fr/images/delete-printer.gif
share/doc/cups-1.1.18.0/fr/images/hold-job.gif
share/doc/cups-1.1.18.0/fr/images/manage-classes.gif
share/doc/cups-1.1.18.0/fr/images/manage-jobs.gif
share/doc/cups-1.1.18.0/fr/images/manage-printers.gif
share/doc/cups-1.1.18.0/fr/images/modify-class.gif
share/doc/cups-1.1.18.0/fr/images/modify-printer.gif
share/doc/cups-1.1.18.0/fr/images/navbar.gif
share/doc/cups-1.1.18.0/fr/images/print-test-page.gif
share/doc/cups-1.1.18.0/fr/images/reject-jobs.gif
share/doc/cups-1.1.18.0/fr/images/release-job.gif
share/doc/cups-1.1.18.0/fr/images/restart-job.gif
share/doc/cups-1.1.18.0/fr/images/show-active.gif
share/doc/cups-1.1.18.0/fr/images/show-completed.gif
share/doc/cups-1.1.18.0/fr/images/start-class.gif
share/doc/cups-1.1.18.0/fr/images/start-printer.gif
share/doc/cups-1.1.18.0/fr/images/stop-class.gif
share/doc/cups-1.1.18.0/fr/images/stop-printer.gif
share/doc/cups-1.1.18.0/fr/index.html
share/doc/cups-1.1.18.0/fr/overview.html
share/doc/cups-1.1.18.0/fr/overview.pdf
share/doc/cups-1.1.18.0/fr/sam.html
share/doc/cups-1.1.18.0/fr/sam.pdf
share/doc/cups-1.1.18.0/fr/sum.html
share/doc/cups-1.1.18.0/fr/sum.pdf
share/doc/cups-1.1.18.0/idd.html
share/doc/cups-1.1.18.0/idd.pdf
share/doc/cups-1.1.18.0/images/accept-jobs.gif
share/doc/cups-1.1.18.0/images/add-class.gif
share/doc/cups-1.1.18.0/images/add-printer.gif
share/doc/cups-1.1.18.0/images/cancel-job.gif
share/doc/cups-1.1.18.0/images/cancel-jobs.gif
share/doc/cups-1.1.18.0/images/cancel.gif
share/doc/cups-1.1.18.0/images/classes.gif
share/doc/cups-1.1.18.0/images/config-printer.gif
share/doc/cups-1.1.18.0/images/continue.gif
share/doc/cups-1.1.18.0/images/cups-block-diagram.gif
share/doc/cups-1.1.18.0/images/cups-large.gif
share/doc/cups-1.1.18.0/images/cups-medium.gif
share/doc/cups-1.1.18.0/images/cups-small.gif
share/doc/cups-1.1.18.0/images/delete-class.gif
share/doc/cups-1.1.18.0/images/delete-printer.gif
share/doc/cups-1.1.18.0/images/draft.gif
share/doc/cups-1.1.18.0/images/hold-job.gif
share/doc/cups-1.1.18.0/images/left.gif
share/doc/cups-1.1.18.0/images/logo.gif
share/doc/cups-1.1.18.0/images/manage-classes.gif
share/doc/cups-1.1.18.0/images/manage-jobs.gif
share/doc/cups-1.1.18.0/images/manage-printers.gif
share/doc/cups-1.1.18.0/images/modify-class.gif
share/doc/cups-1.1.18.0/images/modify-printer.gif
share/doc/cups-1.1.18.0/images/navbar.gif
share/doc/cups-1.1.18.0/images/print-test-page.gif
share/doc/cups-1.1.18.0/images/printer-idle.gif
share/doc/cups-1.1.18.0/images/printer-processing.gif
share/doc/cups-1.1.18.0/images/printer-stopped.gif
share/doc/cups-1.1.18.0/images/reject-jobs.gif
share/doc/cups-1.1.18.0/images/release-job.gif
share/doc/cups-1.1.18.0/images/restart-job.gif
share/doc/cups-1.1.18.0/images/right.gif
share/doc/cups-1.1.18.0/images/show-active.gif
share/doc/cups-1.1.18.0/images/show-completed.gif
share/doc/cups-1.1.18.0/images/start-class.gif
share/doc/cups-1.1.18.0/images/start-printer.gif
share/doc/cups-1.1.18.0/images/stop-class.gif
share/doc/cups-1.1.18.0/images/stop-printer.gif
share/doc/cups-1.1.18.0/index.html
share/doc/cups-1.1.18.0/ipp.html
share/doc/cups-1.1.18.0/ipp.pdf
share/doc/cups-1.1.18.0/overview.html
share/doc/cups-1.1.18.0/overview.pdf
share/doc/cups-1.1.18.0/sam.html
share/doc/cups-1.1.18.0/sam.pdf
share/doc/cups-1.1.18.0/sdd.html
share/doc/cups-1.1.18.0/sdd.pdf
share/doc/cups-1.1.18.0/spm.html
share/doc/cups-1.1.18.0/spm.pdf
share/doc/cups-1.1.18.0/sps.html
share/doc/cups-1.1.18.0/sps.pdf
share/doc/cups-1.1.18.0/ssr.html
share/doc/cups-1.1.18.0/ssr.pdf
share/doc/cups-1.1.18.0/stp.html
share/doc/cups-1.1.18.0/stp.pdf
share/doc/cups-1.1.18.0/sum.html
share/doc/cups-1.1.18.0/sum.pdf
share/doc/cups-1.1.18.0/svd.html
share/doc/cups-1.1.18.0/svd.pdf
share/doc/cups-1.1.18.0/translation.html
share/doc/cups-1.1.18.0/translation.pdf
share/locale/C/cups_C
share/locale/be/cups_be
share/locale/cs/cups_cs
share/locale/de/cups_de
share/locale/en/cups_en
share/locale/es/cups_es
share/locale/fr/cups_fr
share/locale/he/cups_he
share/locale/it/cups_it
share/locale/ru_RU.cp1251/cups_ru_RU.cp1251
share/locale/ru_RU.koi8r/cups_ru_RU.koi8r
share/locale/sv/cups_sv
share/locale/uk/cups_uk
share/locale/uk_UA.cp1251/cups_uk_UA.cp1251
share/locale/zh_CN/cups_zh_CN
@dirrm share/locale/zh_CN
@dirrm share/locale/uk_UA.cp1251
@dirrm share/locale/uk
@dirrm share/locale/sv
@dirrm share/locale/ru_RU.koi8r
@dirrm share/locale/ru_RU.cp1251
@dirrm share/locale/it
@dirrm share/locale/he
@dirrm share/locale/fr
@dirrm share/locale/es
@dirrm share/locale/en
@dirrm share/locale/de
@dirrm share/locale/cs
@dirrm share/locale/be
@dirrm share/locale/C
@dirrm share/doc/cups-1.1.18.0/images
@dirrm share/doc/cups-1.1.18.0/fr/images
@dirrm share/doc/cups-1.1.18.0/fr
@dirrm share/doc/cups-1.1.18.0
@dirrm share/doc
@dirrm share/cups/templates/fr
@dirrm share/cups/templates
@dirrm share/cups/model
@dirrm share/cups/fonts
@dirrm share/cups/data
@dirrm share/cups/charsets
@dirrm share/cups/banners
@dirrm share/cups
@dirrm libexec/cups/filter
@dirrm libexec/cups/daemon
@dirrm libexec/cups/cgi-bin
@dirrm libexec/cups/backend
@dirrm libexec/cups
@dirrm include/cups
@dirrm etc/cups/interfaces
@dirrm etc/cups/certs
@dirrm etc/cups
