# net/radiusd-cistron/pkg-plist
bin/radclient
bin/radlast
bin/radtest
bin/radwho
bin/radzap
etc/raddb/clients.sample
etc/raddb/dictionary
etc/raddb/hints.sample
etc/raddb/huntgroups.sample
etc/raddb/naslist.sample
etc/raddb/naspasswd.sample
etc/raddb/realms.sample
etc/raddb/users.sample
etc/rc.d/radiusd-cistron.sh
sbin/checkrad
sbin/radiusd
sbin/radrelay
share/radius/dictionary.ascend
share/radius/dictionary.bay
share/radius/dictionary.cisco
share/radius/dictionary.cistron
share/radius/dictionary.compat
share/radius/dictionary.erx
share/radius/dictionary.livingston
share/radius/dictionary.redback
share/radius/dictionary.shiva
share/radius/dictionary.standard
share/radius/dictionary.tunnel
share/radius/dictionary.usr
share/radius/dictionary.versanet
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README.Y2K
%%PORTDOCS%%%%DOCSDIR%%/README.cisco
%%PORTDOCS%%%%DOCSDIR%%/README.pam
%%PORTDOCS%%%%DOCSDIR%%/README.proxy
%%PORTDOCS%%%%DOCSDIR%%/README.simul
%%PORTDOCS%%%%DOCSDIR%%/README.usersfile
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm etc/raddb
@dirrm share/radius
@unexec rm -f /var/log/radius.log
@unexec rm -f /var/run/radiusd.pid
@unexec echo "Leaving /var/run/radutmp, /var/log/radwtmp and /var/log/radacct alone."
