# net/miredo/files/isatapd.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: isatapd
# REQUIRE: NETWORKING
# KEYWORD: FreeBSD shutdown

#
# Add the following lines to /etc/rc.conf to enable isatapd:
#
#isatapd_enable="YES"
#
isatapd_enable="${isatapd_enable-NO}"

. /etc/rc.subr

name=isatapd
rcvar=`set_rcvar`

command=%%PREFIX%%/sbin/${name}
required_files=%%PREFIX%%/etc/isatapd.conf

load_rc_config ${name}
run_rc_command "$1"
