# net/miredo/files/miredo-server.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: miredo-server
# REQUIRE: NETWORKING
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable miredo-server:
#
#miredo_server_enable="YES"
#
miredo_server_enable="${miredo_server_enable:-NO}"

. /etc/rc.subr

name=miredo-server
rcvar=miredo_server_enable

command=%%PREFIX%%/sbin/${name}
required_files=%%PREFIX%%/etc/miredo/miredo-server.conf

load_rc_config ${name}
run_rc_command "$1"
