# net/miredo/files/openrc-miredo-server.in
#!/sbin/openrc-run

name=miredo-server

command=%%PREFIX%%/sbin/${name}
required_files=%%PREFIX%%/etc/miredo/miredo-server.conf

depend() {
	use net
	keyword -shutdown
}
