# ports-mgmt/fbsdmon/files/openrc-fbsdmon.in
#!/sbin/openrc-run

name=fbsdmon

command="%%PREFIX%%/sbin/$name"
command_args="-D"

depend() {
		keyword -shutdown
}
