# audio/mt-daapd/pkg-plist
@comment $FreeBSD$
etc/mt-daapd.conf.sample
etc/rc.d/mt-daapd.sh
sbin/mt-daapd
share/mt-daapd/admin-root/aspl-license.html
share/mt-daapd/admin-root/aspl-license.txt
share/mt-daapd/admin-root/config-update.html
share/mt-daapd/admin-root/config.html
share/mt-daapd/admin-root/feedback.html
share/mt-daapd/admin-root/ftr.html
share/mt-daapd/admin-root/gpl-license.html
share/mt-daapd/admin-root/gpl-license.txt
share/mt-daapd/admin-root/hdr.html
share/mt-daapd/admin-root/index.html
share/mt-daapd/admin-root/linkOpaque.gif
share/mt-daapd/admin-root/linkTransparent.gif
share/mt-daapd/admin-root/mt-daapd.css
share/mt-daapd/admin-root/mt-daapd.png
share/mt-daapd/admin-root/required.gif
share/mt-daapd/admin-root/status.html
share/mt-daapd/admin-root/thanks.html
@dirrm var/mt-daapd/db
@dirrm var/mt-daapd
@dirrm share/mt-daapd/mp3
@dirrm share/mt-daapd/admin-root
@dirrm share/mt-daapd
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/CREDITS
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%@dirrm %%DOCSDIR%%
