# sysutils/symon/files/openrc-symux.in
#!/sbin/openrc-run

name="symux"
command="%%PREFIX%%/bin/symux"
