# sysutils/symon/files/openrc-symon.in
#!/sbin/openrc-run

name="symon"
command="%%PREFIX%%/bin/symon"
