# sysutils/symon/files/symux.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: symux
# REQUIRE: NETWORKING SERVERS

. /etc/rc.subr

name="symux"
rcvar=symux_enable
command="%%PREFIX%%/bin/symux"
symux_enable=${symux_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"

