# www/apache20/pkg-plist
bin/htpasswd
bin/htdigest
etc/rc.d/apache.sh
@unexec if cmp -s %D/etc/apache/mime.types %D/etc/apache/mime.types.default; then rm -f %D/etc/apache/mime.types; fi
etc/apache/mime.types.default
@exec [ ! -f %B/mime.types ] && cp %B/%f %B/mime.types
@unexec if cmp -s %D/etc/apache/magic %D/etc/apache/magic.default; then rm -f %D/etc/apache/magic; fi
etc/apache/magic.default
@exec [ ! -f %B/magic ] && cp %B/%f %B/magic
@unexec if cmp -s %D/etc/apache/highperformance-std.conf %D/etc/apache/highperformance-std.conf.default; then rm -f %D/etc/apache/highperformance-std.conf; fi
etc/apache/highperformance-std.conf.default
@exec [ ! -f %B/highperformance-std.conf ] && cp %B/%f %B/highperformance-std.conf
@unexec if cmp -s %D/etc/apache/highperformance.conf %D/etc/apache/highperformance.conf.default; then rm -f %D/etc/apache/highperformance.conf; fi
etc/apache/highperformance.conf.default
@exec [ ! -f %B/highperformance.conf ] && cp %B/%f %B/highperformance.conf
@unexec if cmp -s %D/etc/apache/httpd-std.conf %D/etc/apache/httpd-std.conf.default; then rm -f %D/etc/apache/httpd-std.conf; fi
etc/apache/httpd-std.conf.default
@exec [ ! -f %B/httpd-std.conf ] && cp %B/%f %B/httpd-std.conf
@unexec if cmp -s %D/etc/apache/httpd.conf %D/etc/apache/httpd.conf.default; then rm -f %D/etc/apache/httpd.conf; fi
etc/apache/httpd.conf.default
@exec [ ! -f %B/httpd.conf ] && cp %B/%f %B/httpd.conf
include/apache/apr.h
include/apache/apr_want.h
include/apache/apr_uuid.h
include/apache/apr_user.h
include/apache/apr_time.h
include/apache/apr_thread_proc.h
include/apache/apr_tables.h
include/apache/apr_strings.h
include/apache/apr_signal.h
include/apache/apr_shmem.h
include/apache/apr_portable.h
include/apache/apr_pools.h
include/apache/apr_network_io.h
include/apache/apr_mmap.h
include/apache/apr_md5.h
include/apache/apr_lock.h
include/apache/apr_lib.h
include/apache/apr_hash.h
include/apache/apr_getopt.h
include/apache/apr_general.h
include/apache/apr_fnmatch.h
include/apache/apr_file_io.h
include/apache/apr_file_info.h
include/apache/apr_errno.h
include/apache/apr_dso.h
include/apache/apr_compat.h
include/apache/apr_xlate.h
include/apache/apr_base64.h
include/apache/apu.h
include/apache/apr_xml.h
include/apache/apr_sha1.h
include/apache/apr_sdbm.h
include/apache/apr_ring.h
include/apache/apr_optional.h
include/apache/apr_hooks.h
include/apache/apr_generic_hook.h
include/apache/apr_dbm.h
include/apache/apr_buckets.h
include/apache/apu_compat.h
include/apache/ap_compat.h
include/apache/util_uri.h
include/apache/util_script.h
include/apache/util_md5.h
include/apache/util_filter.h
include/apache/util_ebcdic.h
include/apache/util_date.h
include/apache/util_charset.h
include/apache/util_cfgtree.h
include/apache/scoreboard.h
include/apache/rfc1413.h
include/apache/pcreposix.h
include/apache/mpm_common.h
include/apache/httpd.h
include/apache/http_vhost.h
include/apache/http_request.h
include/apache/http_protocol.h
include/apache/http_main.h
include/apache/http_log.h
include/apache/http_core.h
include/apache/http_connection.h
include/apache/http_config.h
include/apache/ap_release.h
include/apache/ap_mpm.h
include/apache/ap_mmn.h
include/apache/ap_listen.h
include/apache/ap_config_auto.h
include/apache/ap_config.h
include/apache/util_xml.h
include/apache/os.h
include/apache/os-inline.c
include/apache/mpm_default.h
include/apache/mpm.h
include/apache/mod_dav.h
include/apache/mod_include.h
include/apache/pcre.h
include/apache/internal.h
include/apache/config.h
include/apache/unixd.h
lib/libapr.la
lib/libapr.a
lib/libaprutil.la
lib/libaprutil.a
libexec/apache/mod_auth_anon.so
libexec/apache/mod_auth_anon.la
libexec/apache/mod_auth_db.so
libexec/apache/mod_auth_db.la
libexec/apache/mod_auth_digest.so
libexec/apache/mod_auth_digest.la
libexec/apache/mod_file_cache.so
libexec/apache/mod_file_cache.la
libexec/apache/mod_dav_fs.so
libexec/apache/mod_dav_fs.la
libexec/apache/mod_dav.so
libexec/apache/mod_dav.la
libexec/apache/mod_echo.so
libexec/apache/mod_echo.la
libexec/apache/mod_optional_fn_import.so
libexec/apache/mod_optional_fn_import.la
libexec/apache/mod_optional_fn_export.so
libexec/apache/mod_optional_fn_export.la
libexec/apache/mod_include.so
libexec/apache/mod_include.la
libexec/apache/mod_imap.so
libexec/apache/mod_imap.la
libexec/apache/mod_speling.so
libexec/apache/mod_speling.la
libexec/apache/mod_rewrite.so
libexec/apache/mod_rewrite.la
libexec/apache/mod_mime_magic.so
libexec/apache/mod_mime_magic.la
libexec/apache/mod_cern_meta.so
libexec/apache/mod_cern_meta.la
libexec/apache/mod_expires.so
libexec/apache/mod_expires.la
libexec/apache/mod_headers.so
libexec/apache/mod_headers.la
libexec/apache/mod_usertrack.so
libexec/apache/mod_usertrack.la
libexec/apache/mod_unique_id.so
libexec/apache/mod_unique_id.la
libexec/apache/mod_status.so
libexec/apache/mod_status.la
libexec/apache/mod_info.so
libexec/apache/mod_info.la
libexec/apache/mod_cgid.so
libexec/apache/mod_cgid.la
libexec/apache/httpd.exp
sbin/rotatelogs
sbin/htpasswd
sbin/logresolve
sbin/ab
sbin/apachectl
sbin/apxs
sbin/httpd
share/doc/apache/README.rus
@exec mkdir -p %D/www
@exec ln -fs %B %D/www/data.default
@unexec rm -f %D/www/data.default
@exec [ ! -d %D/www/data ] && ln -fs %B %D/www/data
@exec [ ! -d %D/www/cgi-bin ] && ln -fs %D/www/cgi-bin.default %D/www/cgi-bin
share/doc/apache/index.html.tw
share/doc/apache/index.html.se
share/doc/apache/index.html.ru.utf8
share/doc/apache/index.html.ru.ucs4
share/doc/apache/index.html.ru.ucs2
share/doc/apache/index.html.ru.koi8-r
share/doc/apache/index.html.ru.iso-ru
share/doc/apache/index.html.ru.cp866
share/doc/apache/index.html.ru.cp-1251
share/doc/apache/index.html.pt-br
share/doc/apache/index.html.pt
share/doc/apache/index.html.po.iso-pl
share/doc/apache/index.html.no
share/doc/apache/index.html.nl
share/doc/apache/index.html.lu
share/doc/apache/index.html.ltz
share/doc/apache/index.html.kr.iso-kr
share/doc/apache/index.html.kr.iso2022-kr
share/doc/apache/index.html.ja.jis
share/doc/apache/index.html.ja.iso2022-jp
share/doc/apache/index.html.it
share/doc/apache/index.html.he.iso8859-8
share/doc/apache/index.html.fr
share/doc/apache/index.html.et
share/doc/apache/index.html.es
share/doc/apache/index.html.en
share/doc/apache/index.html.el
share/doc/apache/index.html.ee
share/doc/apache/index.html.dk
share/doc/apache/index.html.de
share/doc/apache/index.html.cz
share/doc/apache/index.html.ca
share/doc/apache/apache_pb.gif
share/doc/apache/index.html.tw.Big5
share/doc/apache/manual/misc/tutorials.html
share/doc/apache/manual/misc/FAQ.html
share/doc/apache/manual/misc/known_client_problems.html
share/doc/apache/manual/misc/client_block_api.html
share/doc/apache/manual/misc/compat_notes.html
share/doc/apache/manual/misc/custom_errordocs.html
share/doc/apache/manual/misc/descriptors.html
share/doc/apache/manual/misc/fin_wait_2.html
share/doc/apache/manual/misc/footer.html
share/doc/apache/manual/misc/header.html
share/doc/apache/manual/misc/howto.html
share/doc/apache/manual/misc/index.html
share/doc/apache/manual/misc/perf-tuning.html
share/doc/apache/manual/misc/rewriteguide.html
share/doc/apache/manual/misc/security_tips.html
share/doc/apache/manual/misc/API.html
share/doc/apache/manual/developer/modules.html
share/doc/apache/manual/developer/documenting
share/doc/apache/manual/developer/footer.html
share/doc/apache/manual/developer/header.html
share/doc/apache/manual/developer/hooks.html
share/doc/apache/manual/developer/index.html
share/doc/apache/manual/developer/layeredio.html
share/doc/apache/manual/developer/debugging.html
share/doc/apache/manual/platform/windows.html
share/doc/apache/manual/platform/header.html
share/doc/apache/manual/platform/netware.html
share/doc/apache/manual/platform/perf-bsd44.html
share/doc/apache/manual/platform/perf-dec.html
share/doc/apache/manual/platform/perf-hp.html
share/doc/apache/manual/platform/perf.html
share/doc/apache/manual/platform/readme-tpf.html
share/doc/apache/manual/platform/unixware.html
share/doc/apache/manual/platform/win_compiling.html
share/doc/apache/manual/platform/win_service.html
share/doc/apache/manual/platform/footer.html
share/doc/apache/manual/mod/threaded.html
share/doc/apache/manual/mod/directive-dict.html
share/doc/apache/manual/mod/directives.html
share/doc/apache/manual/mod/footer.html
share/doc/apache/manual/mod/header.html
share/doc/apache/manual/mod/index-bytype.html
share/doc/apache/manual/mod/index.html
share/doc/apache/manual/mod/mod_access.html
share/doc/apache/manual/mod/mod_actions.html
share/doc/apache/manual/mod/mod_alias.html
share/doc/apache/manual/mod/mod_asis.html
share/doc/apache/manual/mod/mod_auth.html
share/doc/apache/manual/mod/mod_auth_anon.html
share/doc/apache/manual/mod/mod_auth_db.html
share/doc/apache/manual/mod/mod_auth_dbm.html
share/doc/apache/manual/mod/mod_auth_digest.html
share/doc/apache/manual/mod/mod_autoindex.html
share/doc/apache/manual/mod/mod_cern_meta.html
share/doc/apache/manual/mod/mod_so.html
share/doc/apache/manual/mod/mod_cgi.html
share/doc/apache/manual/mod/mod_charset_lite.html
share/doc/apache/manual/mod/mod_dav.html
share/doc/apache/manual/mod/mod_dir.html
share/doc/apache/manual/mod/mod_env.html
share/doc/apache/manual/mod/mod_example.html
share/doc/apache/manual/mod/mod_expires.html
share/doc/apache/manual/mod/mod_ext_filter.html
share/doc/apache/manual/mod/mod_file_cache.html
share/doc/apache/manual/mod/mod_headers.html
share/doc/apache/manual/mod/mod_imap.html
share/doc/apache/manual/mod/mod_include.html
share/doc/apache/manual/mod/mod_info.html
share/doc/apache/manual/mod/mod_isapi.html
share/doc/apache/manual/mod/mod_log_config.html
share/doc/apache/manual/mod/mod_mime.html
share/doc/apache/manual/mod/mod_mime_magic.html
share/doc/apache/manual/mod/mod_mmap_static.html
share/doc/apache/manual/mod/mod_negotiation.html
share/doc/apache/manual/mod/mod_proxy.html
share/doc/apache/manual/mod/mod_rewrite.html
share/doc/apache/manual/mod/mod_setenvif.html
share/doc/apache/manual/mod/mod_speling.html
share/doc/apache/manual/mod/mod_status.html
share/doc/apache/manual/mod/mod_unique_id.html
share/doc/apache/manual/mod/mod_userdir.html
share/doc/apache/manual/mod/mod_usertrack.html
share/doc/apache/manual/mod/mod_vhost_alias.html
share/doc/apache/manual/mod/module-dict.html
share/doc/apache/manual/mod/mpm_common.html
share/doc/apache/manual/mod/mpm_winnt.html
share/doc/apache/manual/mod/perchild.html
share/doc/apache/manual/mod/prefork.html
share/doc/apache/manual/mod/core.html
share/doc/apache/manual/search/manual-index.cgi
share/doc/apache/manual/programs/suexec.html
share/doc/apache/manual/programs/apachectl.html
share/doc/apache/manual/programs/apxs.html
share/doc/apache/manual/programs/dbmmanage.html
share/doc/apache/manual/programs/footer.html
share/doc/apache/manual/programs/header.html
share/doc/apache/manual/programs/htdigest.html
share/doc/apache/manual/programs/htpasswd.html
share/doc/apache/manual/programs/httpd.html
share/doc/apache/manual/programs/index.html
share/doc/apache/manual/programs/logresolve.html
share/doc/apache/manual/programs/other.html
share/doc/apache/manual/programs/rotatelogs.html
share/doc/apache/manual/programs/ab.html
share/doc/apache/manual/images/sub.gif
share/doc/apache/manual/images/custom_errordocs.gif
share/doc/apache/manual/images/home.gif
share/doc/apache/manual/images/index.gif
share/doc/apache/manual/images/mod_rewrite_fig1.fig
share/doc/apache/manual/images/mod_rewrite_fig1.gif
share/doc/apache/manual/images/mod_rewrite_fig2.fig
share/doc/apache/manual/images/mod_rewrite_fig2.gif
share/doc/apache/manual/images/pixel.gif
share/doc/apache/manual/images/apache_header.gif
share/doc/apache/manual/howto/cgi.html.en
share/doc/apache/manual/howto/cgi.html.ja.jis
share/doc/apache/manual/howto/footer.html
share/doc/apache/manual/howto/header.html
share/doc/apache/manual/howto/ssi.html.en
share/doc/apache/manual/howto/ssi.html.ja.jis
share/doc/apache/manual/vhosts/virtual-host.html
share/doc/apache/manual/vhosts/examples.html
share/doc/apache/manual/vhosts/fd-limits.html
share/doc/apache/manual/vhosts/footer.html
share/doc/apache/manual/vhosts/header.html
share/doc/apache/manual/vhosts/host.html
share/doc/apache/manual/vhosts/index.html
share/doc/apache/manual/vhosts/ip-based.html
share/doc/apache/manual/vhosts/mass.html
share/doc/apache/manual/vhosts/name-based.html
share/doc/apache/manual/vhosts/vhosts-in-depth.html
share/doc/apache/manual/vhosts/details.html
share/doc/apache/manual/new_features_2_0.html.fr
share/doc/apache/manual/new_features_2_0.html
share/doc/apache/manual/mpm.html
share/doc/apache/manual/location.html
share/doc/apache/manual/invoking.html
share/doc/apache/manual/install.html
share/doc/apache/manual/install-tpf.html
share/doc/apache/manual/index.html.fr
share/doc/apache/manual/index.html
share/doc/apache/manual/header.html
share/doc/apache/manual/handler.html
share/doc/apache/manual/footer.html
share/doc/apache/manual/filter.html
share/doc/apache/manual/env.html
share/doc/apache/manual/ebcdic.html
share/doc/apache/manual/dso.html
share/doc/apache/manual/dns-caveats.html
share/doc/apache/manual/custom-error.html
share/doc/apache/manual/content-negotiation.html
share/doc/apache/manual/configuring.html
share/doc/apache/manual/cgi_path.html
share/doc/apache/manual/bind.html
share/doc/apache/manual/STATUS
share/doc/apache/manual/LICENSE
share/doc/apache/manual/urlmapping.html
share/doc/apache/manual/upgrading.html.fr
share/doc/apache/manual/upgrading.html
share/doc/apache/manual/suexec.html
share/doc/apache/manual/stopping.html
share/doc/apache/manual/server-wide.html
share/doc/apache/manual/sections.html
share/doc/apache/manual/process-model.html
www/cgi-bin.default/test-cgi
www/cgi-bin.default/printenv
www/icons/small/uu.gif
www/icons/small/back.gif
www/icons/small/binary.gif
www/icons/small/binhex.gif
www/icons/small/blank.gif
www/icons/small/broken.gif
www/icons/small/burst.gif
www/icons/small/comp1.gif
www/icons/small/comp2.gif
www/icons/small/compressed.gif
www/icons/small/continued.gif
www/icons/small/dir.gif
www/icons/small/dir2.gif
www/icons/small/doc.gif
www/icons/small/forward.gif
www/icons/small/generic.gif
www/icons/small/generic2.gif
www/icons/small/generic3.gif
www/icons/small/image.gif
www/icons/small/image2.gif
www/icons/small/index.gif
www/icons/small/key.gif
www/icons/small/movie.gif
www/icons/small/patch.gif
www/icons/small/ps.gif
www/icons/small/rainbow.gif
www/icons/small/sound.gif
www/icons/small/sound2.gif
www/icons/small/tar.gif
www/icons/small/text.gif
www/icons/small/transfer.gif
www/icons/small/unknown.gif
www/icons/small/README.txt
www/icons/world1.gif
www/icons/uuencoded.gif
www/icons/uu.gif
www/icons/up.gif
www/icons/unknown.gif
www/icons/transfer.gif
www/icons/text.gif
www/icons/tex.gif
www/icons/tar.gif
www/icons/sphere2.gif
www/icons/sphere1.gif
www/icons/sound2.gif
www/icons/sound1.gif
www/icons/script.gif
www/icons/screw2.gif
www/icons/screw1.gif
www/icons/right.gif
www/icons/quill.gif
www/icons/ps.gif
www/icons/portal.gif
www/icons/pie8.gif
www/icons/pie7.gif
www/icons/pie6.gif
www/icons/pie5.gif
www/icons/pie4.gif
www/icons/pie3.gif
www/icons/pie2.gif
www/icons/pie1.gif
www/icons/pie0.gif
www/icons/pdf.gif
www/icons/patch.gif
www/icons/p.gif
www/icons/movie.gif
www/icons/link.gif
www/icons/left.gif
www/icons/layout.gif
www/icons/index.gif
www/icons/image3.gif
www/icons/image2.gif
www/icons/image1.gif
www/icons/icon.sheet.gif
www/icons/hand.up.gif
www/icons/hand.right.gif
www/icons/generic.sec.gif
www/icons/generic.red.gif
www/icons/generic.gif
www/icons/forward.gif
www/icons/folder.sec.gif
www/icons/folder.gif
www/icons/folder.open.gif
www/icons/f.gif
www/icons/dvi.gif
www/icons/down.gif
www/icons/dir.gif
www/icons/continued.gif
www/icons/compressed.gif
www/icons/comp.gray.gif
www/icons/comp.blue.gif
www/icons/c.gif
www/icons/burst.gif
www/icons/broken.gif
www/icons/box2.gif
www/icons/box1.gif
www/icons/bomb.gif
www/icons/blank.gif
www/icons/binhex.gif
www/icons/binary.gif
www/icons/ball.red.gif
www/icons/ball.gray.gif
www/icons/back.gif
www/icons/apache_pb.gif
www/icons/alert.red.gif
www/icons/alert.black.gif
www/icons/a.gif
www/icons/world2.gif
www/icons/README
@unexec if [ -L %D/www/cgi-bin ]; then rm -f %D/www/cgi-bin; fi
@unexec if [ -L %D/www/data ]; then rm -f %D/www/data; fi
@dirrm etc/apache
@dirrm include/apache
@dirrm libexec/apache
@dirrm share/doc/apache/manual/vhosts
@dirrm share/doc/apache/manual/search
@dirrm share/doc/apache/manual/programs
@dirrm share/doc/apache/manual/platform
@dirrm share/doc/apache/manual/mod
@dirrm share/doc/apache/manual/misc
@dirrm share/doc/apache/manual/images
@dirrm share/doc/apache/manual/howto
@dirrm share/doc/apache/manual/developer
@dirrm share/doc/apache/manual
@dirrm share/doc/apache
@dirrm www/icons/small
@dirrm www/icons
@dirrm www/cgi-bin.default
@dirrm www
