# sysutils/vm-bhyve/pkg-plist
etc/rc.d/vm
lib/vm-bhyve/vm-cmd
lib/vm-bhyve/vm-common
lib/vm-bhyve/vm-config
lib/vm-bhyve/vm-core
lib/vm-bhyve/vm-guest
lib/vm-bhyve/vm-run
lib/vm-bhyve/vm-switch
lib/vm-bhyve/vm-sysrc
lib/vm-bhyve/vm-zfs
man/man8/vm.8.gz
sbin/vm
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/alpine.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/centos.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/debian.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/default.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/freebsd-zvol.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/netbsd.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/openbsd.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/ubuntu.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/windows.conf
