# sysutils/vm-bhyve/pkg-plist
etc/rc.d/vm
lib/vm-bhyve/vm-base
lib/vm-bhyve/vm-cmd
lib/vm-bhyve/vm-config
lib/vm-bhyve/vm-core
lib/vm-bhyve/vm-datastore
lib/vm-bhyve/vm-guest
lib/vm-bhyve/vm-info
lib/vm-bhyve/vm-migration
lib/vm-bhyve/vm-rctl
lib/vm-bhyve/vm-run
lib/vm-bhyve/vm-switch
lib/vm-bhyve/vm-switch-manual
lib/vm-bhyve/vm-switch-standard
lib/vm-bhyve/vm-switch-vale
lib/vm-bhyve/vm-switch-vxlan
lib/vm-bhyve/vm-util
lib/vm-bhyve/vm-zfs
man/man8/vm.8.gz
sbin/vm
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/alpine.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/arch.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/centos6.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/centos7.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/config.sample
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/coreos.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/debian.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/default.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/dragonfly.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/freebsd-zvol.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/freepbx.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/linux-zvol.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/netbsd.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/openbsd.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/resflash.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/ubuntu.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/windows.conf
