# www/resin2/pkg-plist
resin/doc/beta.xtp
resin/doc/changes-1.xtp
resin/doc/changes.xtp
resin/doc/errata.xtp
resin/doc/features.xtp
resin/doc/future.xtp
resin/doc/index.xtp
resin/doc/toc.xml
resin/doc/examples/basic/WEB-INF/classes/Env.java
resin/doc/examples/basic/WEB-INF/classes/HelloServlet.java
resin/doc/examples/basic/counter.jsp
resin/doc/examples/basic/dir.jsp
resin/doc/examples/basic/dummy.jsp
resin/doc/examples/basic/env.jsp
resin/doc/examples/basic/file.jsp
resin/doc/examples/basic/form.html
resin/doc/examples/basic/form.jsp
resin/doc/examples/basic/hello.jsp
resin/doc/examples/basic/quote.jsp
resin/doc/examples/basic/quote_error.jsp
resin/doc/examples/basic/session.jsp
resin/doc/examples/basic/stats.jsp
resin/doc/examples/basic/com/caucho/fortune/1.html
resin/doc/examples/basic/com/caucho/fortune/2.html
resin/doc/examples/basic/com/caucho/fortune/3.html
resin/doc/examples/basic/com/caucho/fortune/4.html
resin/doc/examples/basic/com/caucho/fortune/5.html
resin/doc/examples/basic/data/fortune/1.html
resin/doc/examples/basic/data/fortune/2.html
resin/doc/examples/basic/data/fortune/3.html
resin/doc/examples/basic/data/fortune/4.html
resin/doc/examples/basic/data/fortune/5.html
resin/doc/examples/index.xtp
resin/doc/examples/toc.xml
resin/doc/examples/view_source.jsp
resin/doc/examples/login/WEB-INF/classes/test/TestAuthenticator.java
resin/doc/examples/login/WEB-INF/web.xml
resin/doc/examples/login/error.html
resin/doc/examples/login/index.jsp
resin/doc/examples/login/login.html
resin/doc/examples/login/login.jsp
resin/doc/examples/login/protected.jsp
resin/doc/examples/navigation/WEB-INF/scripts/Navigation.js
resin/doc/examples/navigation/WEB-INF/scripts/fancy_header.js
resin/doc/examples/navigation/WEB-INF/scripts/header.js
resin/doc/examples/navigation/article1.jsp
resin/doc/examples/navigation/article2.jsp
resin/doc/examples/navigation/child1.jsp
resin/doc/examples/navigation/child2.jsp
resin/doc/examples/navigation/conclusion.jsp
resin/doc/examples/navigation/navigation.jsp
resin/doc/examples/navigation/toc.xml
resin/doc/examples/templates/admin.jsp
resin/doc/examples/templates/add.jsp
resin/doc/examples/templates/add.xtp
resin/doc/examples/templates/WEB-INF/classes/jsp/Guest.java
resin/doc/examples/templates/WEB-INF/classes/jsp/GuestBook.java
resin/doc/examples/templates/WEB-INF/classes/jsp/GuestJsp.java
resin/doc/examples/templates/WEB-INF/classes/jsp/GuestXtp.java
resin/doc/examples/templates/WEB-INF/classes/jsp/RSS.java
resin/doc/examples/templates/WEB-INF/xsl/caucho_form.xsl
resin/doc/examples/templates/WEB-INF/xsl/default.xsl
resin/doc/examples/templates/WEB-INF/xsl/defaultcopy.xsl
resin/doc/examples/templates/WEB-INF/xsl/form.xsl
resin/doc/examples/templates/WEB-INF/xsl/formatting.xsl
resin/doc/examples/templates/WEB-INF/xsl/guestbook.xsl
resin/doc/examples/templates/WEB-INF/xsl/header.xsl
resin/doc/examples/templates/WEB-INF/xsl/rss.xsl
resin/doc/examples/templates/admin.xtp
resin/doc/examples/templates/guestbook.jsp
resin/doc/examples/templates/guestbook.xtp
resin/doc/examples/templates/guestbook_bean.jsp
resin/doc/examples/templates/guestbook_inc.jsp
resin/doc/examples/templates/list.jsp
resin/doc/examples/templates/login.jsp
resin/doc/examples/templates/login.xtp
resin/doc/examples/templates/loginold.jsp
resin/doc/examples/templates/rss.xtp
resin/doc/examples/templates/simple-format.xsl
resin/doc/examples/tictactoe/WEB-INF/scripts/score.js
resin/doc/examples/tictactoe/error.jsp
resin/doc/examples/tictactoe/tictactoe.jsp
resin/doc/examples/tutorial/WEB-INF/classes/servlet/trivial/HelloServlet.java
resin/doc/examples/tutorial/WEB-INF/web.xml
resin/doc/examples/xsl/WEB-INF/classes/xsl/Dispatch.java
resin/doc/examples/xsl/WEB-INF/web.xml
resin/doc/examples/xsl/WEB-INF/scripts/Navigation.js
resin/doc/examples/xsl/WEB-INF/xsl/basic.xsl
resin/doc/examples/xsl/WEB-INF/xsl/content.xsl
resin/doc/examples/xsl/WEB-INF/xsl/copy.xsl
resin/doc/examples/xsl/WEB-INF/xsl/counter.xsl
resin/doc/examples/xsl/WEB-INF/xsl/default.xsl
resin/doc/examples/xsl/WEB-INF/xsl/header.xsl
resin/doc/examples/xsl/WEB-INF/xsl/plain.xsl
resin/doc/examples/xsl/WEB-INF/xsl/verify.xsl
resin/doc/examples/xsl/basic.xtp
resin/doc/examples/xsl/counter.xtp
resin/doc/examples/xsl/good-doc.xtp
resin/doc/examples/xsl/index.xml
resin/doc/examples/xsl/navigation.xtp
resin/doc/examples/xsl/toc.xml
resin/doc/examples/xsl/vary.xml
resin/doc/examples/xsl/vary.xtp
resin/doc/examples/xsl/verify-bad.xtp
resin/doc/examples/xsl/verify-ok.xtp
resin/doc/examples/xsl/data/fortune/1.html
resin/doc/examples/xsl/data/fortune/2.html
resin/doc/examples/xsl/data/fortune/3.html
resin/doc/examples/xsl/data/fortune/4.html
resin/doc/examples/xsl/data/fortune/5.html
resin/doc/images/apache_srun.gif
resin/doc/images/appserver.gif
resin/doc/images/background.gif
resin/doc/images/backup.gif
resin/doc/images/ball8.gif
resin/doc/images/bean.gif
resin/doc/images/bean1.gif
resin/doc/images/caucho.gif
resin/doc/images/compiler.gif
resin/doc/images/file.gif
resin/doc/images/folder.gif
resin/doc/images/hbleed.gif
resin/doc/images/httpd.gif
resin/doc/images/iis_srun.gif
resin/doc/images/jndi_tree.gif
resin/doc/images/jsp.gif
resin/doc/images/logo.gif
resin/doc/images/pixel.gif
resin/doc/images/resin_powered.gif
resin/doc/images/same_srun.gif
resin/doc/images/script.gif
resin/doc/images/srunc.gif
resin/doc/images/sticky-sessions.gif
resin/doc/images/style.gif
resin/doc/images/webserver.gif
resin/doc/images/website_srun.gif
resin/doc/images/caucho-blue.gif
resin/doc/images/left_background.gif
resin/doc/java_tut/WEB-INF/classes/example/servlet/basic/HelloServlet.java
resin/doc/java_tut/WEB-INF/classes/example/servlet/basic/InitServlet.java
resin/doc/java_tut/WEB-INF/classes/example/servlet/basic/MappingServlet.java
resin/doc/java_tut/WEB-INF/classes/example/servlet/basic/URLServlet.java
resin/doc/java_tut/WEB-INF/classes/example/servlet/database/InitDatabaseServlet.java
resin/doc/java_tut/WEB-INF/classes/example/servlet/database/QueryServlet.java
resin/doc/java_tut/WEB-INF/web.xml
resin/doc/java_tut/WEB-INF/sql/default.sql
resin/doc/java_tut/add_page.xtp
resin/doc/java_tut/bean.xtp
resin/doc/java_tut/cache.xtp
resin/doc/java_tut/conf-basic.xtp
resin/doc/java_tut/conf-dir.xtp
resin/doc/java_tut/conf-servlet.xtp
resin/doc/java_tut/conf.xtp
resin/doc/java_tut/db-config.xtp
resin/doc/java_tut/db_form.xtp
resin/doc/java_tut/debug.xtp
resin/doc/java_tut/filter-vary.xtp
resin/doc/java_tut/form.xtp
resin/doc/java_tut/hardcore.xtp
resin/doc/java_tut/index.xtp
resin/doc/java_tut/isp.xtp
resin/doc/java_tut/linux-boot.xtp
resin/doc/java_tut/linux-init.xtp
resin/doc/java_tut/mail_form.xtp
resin/doc/java_tut/multipart.xtp
resin/doc/java_tut/parse_error.xtp
resin/doc/java_tut/req_headers.xtp
resin/doc/java_tut/plugin-dispatch.xtp
resin/doc/java_tut/req_url.xtp
resin/doc/java_tut/run-at.xtp
resin/doc/java_tut/run_error.xtp
resin/doc/java_tut/servlets.xtp
resin/doc/java_tut/tag-attribute.xtp
resin/doc/java_tut/tag-connection.xtp
resin/doc/java_tut/tag-hello.xtp
resin/doc/java_tut/tag-mail.xtp
resin/doc/java_tut/tag-xsl.xtp
resin/doc/java_tut/tcp-sessions.xtp
resin/doc/java_tut/toc.xml
resin/doc/java_tut/tuning.xtp
resin/doc/java_tut/war-dir.xtp
resin/doc/java_tut/xsl-filter.xtp
resin/doc/java_tut/xtp-copy.xtp
resin/doc/java_tut/xtp-format.xtp
resin/doc/java_tut/xtp-jsp-stylescript.xtp
resin/doc/java_tut/xtp-jsp.xtp
resin/doc/java_tut/xtp-page.xtp
resin/doc/java_tut/xtp-strict.xtp
resin/doc/java_tut/javadoc/constant-values.html
resin/doc/java_tut/javadoc/resources/inherit.gif
resin/doc/java_tut/javadoc/overview-tree.html
resin/doc/java_tut/javadoc/index-all.html
resin/doc/java_tut/javadoc/deprecated-list.html
resin/doc/java_tut/javadoc/allclasses-frame.html
resin/doc/java_tut/javadoc/allclasses-noframe.html
resin/doc/java_tut/javadoc/index.html
resin/doc/java_tut/javadoc/packages.html
resin/doc/java_tut/javadoc/overview-summary.html
resin/doc/java_tut/javadoc/overview-frame.html
resin/doc/java_tut/javadoc/example/servlet/basic/package-frame.html
resin/doc/java_tut/javadoc/example/servlet/basic/package-summary.html
resin/doc/java_tut/javadoc/example/servlet/basic/package-tree.html
resin/doc/java_tut/javadoc/example/servlet/basic/HelloServlet.html
resin/doc/java_tut/javadoc/example/servlet/basic/InitServlet.html
resin/doc/java_tut/javadoc/example/servlet/basic/MappingServlet.html
resin/doc/java_tut/javadoc/example/servlet/basic/URLServlet.html
resin/doc/java_tut/javadoc/example/servlet/database/package-frame.html
resin/doc/java_tut/javadoc/example/servlet/database/package-summary.html
resin/doc/java_tut/javadoc/example/servlet/database/package-tree.html
resin/doc/java_tut/javadoc/example/servlet/database/InitDatabaseServlet.html
resin/doc/java_tut/javadoc/example/servlet/database/QueryServlet.html
resin/doc/java_tut/javadoc/serialized-form.html
resin/doc/java_tut/javadoc/package-list
resin/doc/java_tut/javadoc/help-doc.html
resin/doc/java_tut/javadoc/stylesheet.css
resin/doc/java_tut/images/same_srun.gif
resin/doc/java_tut/images/srunc.gif
resin/doc/ref/jsp/jsp-el.xtp
resin/doc/ref/admin.xtp
resin/doc/ref/app-config.xtp
resin/doc/ref/atp.xtp
resin/doc/ref/auth-config.xtp
resin/doc/ref/balance.xtp
resin/doc/ref/basic-config.xtp
resin/doc/ref/bnf.xtp
resin/doc/ref/burlap.xtp
resin/doc/ref/caching.xtp
resin/doc/ref/changes.xtp
resin/doc/ref/config-index.xtp
resin/doc/ref/config-sum.xtp
resin/doc/ref/config.xtp
resin/doc/ref/contents.xtp
resin/doc/ref/cse-apache-unix.xtp
resin/doc/ref/cse-apache-win.xtp
resin/doc/ref/cse-iis.xtp
resin/doc/ref/cse-netscape.xtp
resin/doc/ref/cse-oreilly-website.xtp
resin/doc/ref/db-config.xtp
resin/doc/ref/xml/attributes.xtp
resin/doc/ref/xml/elements.xtp
resin/doc/ref/xml/get_elements.xtp
resin/doc/ref/xml/html.xtp
resin/doc/ref/xml/index.xtp
resin/doc/ref/xml/intro.xtp
resin/doc/ref/xml/parsing.xtp
resin/doc/ref/xml/sax.xtp
resin/doc/ref/xml/text.xtp
resin/doc/ref/xml/vfs.xtp
resin/doc/ref/ejb-ias.xtp
resin/doc/ref/ejb-jboss.xtp
resin/doc/ref/ejb-jonas.xtp
resin/doc/ref/ejb-orion.xtp
resin/doc/ref/ejb-weblogic.xtp
resin/doc/ref/ejb.xtp
resin/doc/ref/errata.xtp
resin/doc/ref/faq.xtp
resin/doc/ref/file-config.xtp
resin/doc/ref/filesystems.xtp
resin/doc/ref/filter-config.xtp
resin/doc/ref/filter-library.xtp
resin/doc/ref/filter.xtp
resin/doc/ref/glossary.xtp
resin/doc/ref/hessian.xtp
resin/doc/ref/http-config.xtp
resin/doc/ref/httpd.xtp
resin/doc/ref/ide-jbuilder.xtp
resin/doc/ref/ide-kawa.xtp
resin/doc/ref/ide.xtp
resin/doc/ref/index.xtp
resin/doc/ref/javascript.xtp
resin/doc/ref/jms.xtp
resin/doc/ref/jmx.xtp
resin/doc/ref/jndi.xtp
resin/doc/ref/jsp-actions.xtp
resin/doc/ref/jsp-applications.xtp
resin/doc/ref/jsp-directives.xtp
resin/doc/ref/jsp-variables.xtp
resin/doc/ref/jsp.xtp
resin/doc/ref/jspapp.xtp
resin/doc/ref/jtp.xtp
resin/doc/ref/log-config.xtp
resin/doc/ref/macos-x.xtp
resin/doc/ref/page_context.xtp
resin/doc/ref/port-config.xtp
resin/doc/ref/raw-socket.xtp
resin/doc/ref/resin-config.xtp
resin/doc/ref/security.xtp
resin/doc/ref/servlet-config.xtp
resin/doc/ref/servlet-library.xtp
resin/doc/ref/servlet.xtp
resin/doc/ref/servletcontext.xtp
resin/doc/ref/sessions.xtp
resin/doc/ref/xsl/apply-templates.xtp
resin/doc/ref/xsl/attribute.xtp
resin/doc/ref/xsl/boxes.xtp
resin/doc/ref/xsl/command-line.xtp
resin/doc/ref/xsl/copy.xtp
resin/doc/ref/xsl/filter-attribute.xtp
resin/doc/ref/xsl/filter.xtp
resin/doc/ref/xsl/hello.xtp
resin/doc/ref/xsl/if-attribute.xtp
resin/doc/ref/xsl/if-xpath-attribute.xtp
resin/doc/ref/xsl/if.xtp
resin/doc/ref/xsl/index.xtp
resin/doc/ref/xsl/modes.xtp
resin/doc/ref/xsl/ref.xtp
resin/doc/ref/xsl/short-attribute.xtp
resin/doc/ref/xsl/value-of.xtp
resin/doc/ref/xsl/xpath-find.xtp
resin/doc/ref/xsl/xpath-select.xtp
resin/doc/ref/xsl/xpath.xtp
resin/doc/ref/xsl/xsl-api.xtp
resin/doc/ref/xsl/xsl.xtp
resin/doc/ref/xsl/xsl-boxes.xtp
resin/doc/ref/srun-other.xtp
resin/doc/ref/ssl-quick.xml
resin/doc/ref/ssl-quick.xtp
resin/doc/ref/ssl.xml
resin/doc/ref/ssl.xtp
resin/doc/ref/starting.xtp
resin/doc/ref/static.xtp
resin/doc/ref/stylescript.xtp
resin/doc/ref/taglib.xtp
resin/doc/ref/toc.xml
resin/doc/ref/topics.xtp
resin/doc/ref/trouble.xtp
resin/doc/ref/velocity.xtp
resin/doc/ref/virtual-host.xtp
resin/doc/ref/visual-age.xtp
resin/doc/ref/war-dir.xtp
resin/doc/ref/webdav.xtp
resin/doc/ref/xml.xtp
resin/doc/ref/xpath-fun.xtp
resin/doc/ref/xpath.xtp
resin/doc/ref/xsl.xtp
resin/doc/ref/xslt-lite-eg.xtp
resin/doc/ref/xslt-lite.xtp
resin/doc/ref/xtp-jsp.xtp
resin/doc/ref/xtp-index.xtp
resin/doc/ref/xtp-intro.xtp
resin/doc/ref/xtp.xtp
resin/doc/faq/classloader.xtp
resin/doc/faq/config.xtp
resin/doc/faq/debug.xtp
resin/doc/faq/index.xtp
resin/doc/faq/install.xtp
resin/doc/faq/misc.xtp
resin/doc/faq/servlet.xtp
resin/doc/faq/ssl.xtp
resin/doc/faq/toc.xml
resin/doc/css/default.css
resin/lib/resin.jar
resin/lib/jsdk23.jar
resin/lib/jdbc2_0-stdext.jar
resin/lib/jta_101.jar
resin/lib/jndi.jar
resin/lib/dom.jar
resin/lib/sax.jar
resin/lib/jaxp.jar
resin/lib/jmx.jar
resin/lib/webutil.jar
resin/lib/jdbc-mysql.jar
resin/xsl/caucho/comment.xsl
resin/xsl/caucho/contents.xsl
resin/xsl/caucho/control.xsl
resin/xsl/caucho/debug.xsl
resin/xsl/caucho/defaultcopy.xsl
resin/xsl/caucho/format.xsl
resin/xsl/caucho/newstyle.xsl
resin/xsl/caucho/weblog.xsl
resin/xsl/caucho_default.xsl
resin/xsl/caucho_header.xsl
resin/xsl/caucho_subpage.xsl
resin/xsl/default.xsl
resin/xsl/doc_header.xsl
resin/xsl/header.xsl
resin/xsl/js_header.xsl
resin/xsl/new_doc_header.xsl
resin/xsl/search.xsl
resin/xsl/subpage.xsl
resin/xsl/js/caucho/defaultcopy.xsl
resin/xsl/js/caucho/format.xsl
resin/xsl/js/caucho/weblog.xsl
etc/rc.d/resin.sh
sbin/resinctl
%%MOD_DIR%%/mod_caucho.so
@dirrm resin/doc/examples/basic/WEB-INF/classes
@dirrm resin/doc/examples/basic/WEB-INF/tmp
@dirrm resin/doc/examples/basic/WEB-INF
@dirrm resin/doc/examples/basic/com/caucho/fortune
@dirrm resin/doc/examples/basic/com/caucho
@dirrm resin/doc/examples/basic/com
@dirrm resin/doc/examples/basic/data/fortune
@dirrm resin/doc/examples/basic/data
@dirrm resin/doc/examples/basic
@dirrm resin/doc/examples/login/WEB-INF/classes/test
@dirrm resin/doc/examples/login/WEB-INF/classes
@dirrm resin/doc/examples/login/WEB-INF/tmp
@dirrm resin/doc/examples/login/WEB-INF
@dirrm resin/doc/examples/login
@dirrm resin/doc/examples/navigation/WEB-INF/scripts
@dirrm resin/doc/examples/navigation/WEB-INF/tmp
@dirrm resin/doc/examples/navigation/WEB-INF
@dirrm resin/doc/examples/navigation
@dirrm resin/doc/examples/templates/WEB-INF/classes/jsp
@dirrm resin/doc/examples/templates/WEB-INF/classes
@dirrm resin/doc/examples/templates/WEB-INF/xsl
@dirrm resin/doc/examples/templates/WEB-INF/tmp
@dirrm resin/doc/examples/templates/WEB-INF
@dirrm resin/doc/examples/templates
@dirrm resin/doc/examples/tictactoe/WEB-INF/scripts
@dirrm resin/doc/examples/tictactoe/WEB-INF/tmp
@dirrm resin/doc/examples/tictactoe/WEB-INF
@dirrm resin/doc/examples/tictactoe
@dirrm resin/doc/examples/tutorial/WEB-INF/classes/servlet/trivial
@dirrm resin/doc/examples/tutorial/WEB-INF/classes/servlet
@dirrm resin/doc/examples/tutorial/WEB-INF/classes
@dirrm resin/doc/examples/tutorial/WEB-INF
@dirrm resin/doc/examples/tutorial
@dirrm resin/doc/examples/xsl/WEB-INF/classes/xsl
@dirrm resin/doc/examples/xsl/WEB-INF/classes
@dirrm resin/doc/examples/xsl/WEB-INF/scripts
@dirrm resin/doc/examples/xsl/WEB-INF/xsl
@dirrm resin/doc/examples/xsl/WEB-INF/tmp
@dirrm resin/doc/examples/xsl/WEB-INF
@dirrm resin/doc/examples/xsl/data/fortune
@dirrm resin/doc/examples/xsl/data
@dirrm resin/doc/examples/xsl
@dirrm resin/doc/examples/tags/WEB-INF/tmp
@dirrm resin/doc/examples/tags/WEB-INF
@dirrm resin/doc/examples/tags
@dirrm resin/doc/examples
@dirrm resin/doc/images
@dirrm resin/doc/java_tut/WEB-INF/classes/example/servlet/basic
@dirrm resin/doc/java_tut/WEB-INF/classes/example/servlet/database
@dirrm resin/doc/java_tut/WEB-INF/classes/example/servlet
@dirrm resin/doc/java_tut/WEB-INF/classes/example
@dirrm resin/doc/java_tut/WEB-INF/classes
@dirrm resin/doc/java_tut/WEB-INF/sql
@dirrm resin/doc/java_tut/WEB-INF/tmp
@dirrm resin/doc/java_tut/WEB-INF
@dirrm resin/doc/java_tut/javadoc/resources
@dirrm resin/doc/java_tut/javadoc/example/servlet/basic
@dirrm resin/doc/java_tut/javadoc/example/servlet/database
@dirrm resin/doc/java_tut/javadoc/example/servlet
@dirrm resin/doc/java_tut/javadoc/example
@dirrm resin/doc/java_tut/javadoc
@dirrm resin/doc/java_tut/images
@dirrm resin/doc/java_tut
@dirrm resin/doc/ref/jsp
@dirrm resin/doc/ref/xml
@dirrm resin/doc/ref/xsl
@dirrm resin/doc/ref
@dirrm resin/doc/faq
@dirrm resin/doc/css
@dirrm resin/doc/WEB-INF/classes
@dirrm resin/doc/WEB-INF/lib
@dirrm resin/doc/WEB-INF
@dirrm resin/doc
@dirrm resin/lib
@dirrm resin/webapps
@dirrm resin/xsl/caucho
@dirrm resin/xsl/js/caucho
@dirrm resin/xsl/js
@dirrm resin/xsl
@dirrm resin
@unexec if cmp -s %D/etc/resin.xml %D/etc/resin.xml-dist; then rm -f %D/etc/resin.xml; else echo If permanently deleting this package, %%PREFIX%%/etc/resin.xml must be removed manually; fi
etc/resin.xml-dist
