# www/resin2/pkg-plist
@comment $FreeBSD$
%%APP_NAME%%/doc/beta.xtp
%%APP_NAME%%/doc/changes-1.xtp
%%APP_NAME%%/doc/changes.xtp
%%APP_NAME%%/doc/css/default.css
%%APP_NAME%%/doc/errata.xtp
%%APP_NAME%%/doc/examples/basic/WEB-INF/classes/Env.java
%%APP_NAME%%/doc/examples/basic/WEB-INF/classes/HelloServlet.java
%%APP_NAME%%/doc/examples/basic/com/caucho/fortune/1.html
%%APP_NAME%%/doc/examples/basic/com/caucho/fortune/2.html
%%APP_NAME%%/doc/examples/basic/com/caucho/fortune/3.html
%%APP_NAME%%/doc/examples/basic/com/caucho/fortune/4.html
%%APP_NAME%%/doc/examples/basic/com/caucho/fortune/5.html
%%APP_NAME%%/doc/examples/basic/counter.jsp
%%APP_NAME%%/doc/examples/basic/data/fortune/1.html
%%APP_NAME%%/doc/examples/basic/data/fortune/2.html
%%APP_NAME%%/doc/examples/basic/data/fortune/3.html
%%APP_NAME%%/doc/examples/basic/data/fortune/4.html
%%APP_NAME%%/doc/examples/basic/data/fortune/5.html
%%APP_NAME%%/doc/examples/basic/dir.jsp
%%APP_NAME%%/doc/examples/basic/dummy.jsp
%%APP_NAME%%/doc/examples/basic/env.jsp
%%APP_NAME%%/doc/examples/basic/file.jsp
%%APP_NAME%%/doc/examples/basic/form.html
%%APP_NAME%%/doc/examples/basic/form.jsp
%%APP_NAME%%/doc/examples/basic/hello.jsp
%%APP_NAME%%/doc/examples/basic/quote.jsp
%%APP_NAME%%/doc/examples/basic/quote_error.jsp
%%APP_NAME%%/doc/examples/basic/session.jsp
%%APP_NAME%%/doc/examples/basic/stats.jsp
%%APP_NAME%%/doc/examples/index.xtp
%%APP_NAME%%/doc/examples/login/WEB-INF/classes/test/TestAuthenticator.java
%%APP_NAME%%/doc/examples/login/WEB-INF/web.xml
%%APP_NAME%%/doc/examples/login/error.html
%%APP_NAME%%/doc/examples/login/index.jsp
%%APP_NAME%%/doc/examples/login/login.html
%%APP_NAME%%/doc/examples/login/login.jsp
%%APP_NAME%%/doc/examples/login/protected.jsp
%%APP_NAME%%/doc/examples/navigation/WEB-INF/scripts/Navigation.js
%%APP_NAME%%/doc/examples/navigation/WEB-INF/scripts/fancy_header.js
%%APP_NAME%%/doc/examples/navigation/WEB-INF/scripts/header.js
%%APP_NAME%%/doc/examples/navigation/article1.jsp
%%APP_NAME%%/doc/examples/navigation/article2.jsp
%%APP_NAME%%/doc/examples/navigation/child1.jsp
%%APP_NAME%%/doc/examples/navigation/child2.jsp
%%APP_NAME%%/doc/examples/navigation/conclusion.jsp
%%APP_NAME%%/doc/examples/navigation/navigation.jsp
%%APP_NAME%%/doc/examples/navigation/toc.xml
%%APP_NAME%%/doc/examples/templates/WEB-INF/classes/jsp/Guest.java
%%APP_NAME%%/doc/examples/templates/WEB-INF/classes/jsp/GuestBook.java
%%APP_NAME%%/doc/examples/templates/WEB-INF/classes/jsp/GuestJsp.java
%%APP_NAME%%/doc/examples/templates/WEB-INF/classes/jsp/GuestXtp.java
%%APP_NAME%%/doc/examples/templates/WEB-INF/classes/jsp/RSS.java
%%APP_NAME%%/doc/examples/templates/WEB-INF/xsl/caucho_form.xsl
%%APP_NAME%%/doc/examples/templates/WEB-INF/xsl/default.xsl
%%APP_NAME%%/doc/examples/templates/WEB-INF/xsl/defaultcopy.xsl
%%APP_NAME%%/doc/examples/templates/WEB-INF/xsl/form.xsl
%%APP_NAME%%/doc/examples/templates/WEB-INF/xsl/formatting.xsl
%%APP_NAME%%/doc/examples/templates/WEB-INF/xsl/guestbook.xsl
%%APP_NAME%%/doc/examples/templates/WEB-INF/xsl/header.xsl
%%APP_NAME%%/doc/examples/templates/WEB-INF/xsl/rss.xsl
%%APP_NAME%%/doc/examples/templates/add.jsp
%%APP_NAME%%/doc/examples/templates/add.xtp
%%APP_NAME%%/doc/examples/templates/admin.jsp
%%APP_NAME%%/doc/examples/templates/admin.xtp
%%APP_NAME%%/doc/examples/templates/guestbook.jsp
%%APP_NAME%%/doc/examples/templates/guestbook.xtp
%%APP_NAME%%/doc/examples/templates/guestbook_bean.jsp
%%APP_NAME%%/doc/examples/templates/guestbook_inc.jsp
%%APP_NAME%%/doc/examples/templates/list.jsp
%%APP_NAME%%/doc/examples/templates/login.jsp
%%APP_NAME%%/doc/examples/templates/login.xtp
%%APP_NAME%%/doc/examples/templates/loginold.jsp
%%APP_NAME%%/doc/examples/templates/rss.xtp
%%APP_NAME%%/doc/examples/templates/simple-format.xsl
%%APP_NAME%%/doc/examples/tictactoe/WEB-INF/scripts/score.js
%%APP_NAME%%/doc/examples/tictactoe/error.jsp
%%APP_NAME%%/doc/examples/tictactoe/tictactoe.jsp
%%APP_NAME%%/doc/examples/toc.xml
%%APP_NAME%%/doc/examples/tutorial/WEB-INF/classes/servlet/trivial/HelloServlet.java
%%APP_NAME%%/doc/examples/tutorial/WEB-INF/web.xml
%%APP_NAME%%/doc/examples/view_source.jsp
%%APP_NAME%%/doc/examples/xsl/WEB-INF/classes/xsl/Dispatch.java
%%APP_NAME%%/doc/examples/xsl/WEB-INF/scripts/Navigation.js
%%APP_NAME%%/doc/examples/xsl/WEB-INF/web.xml
%%APP_NAME%%/doc/examples/xsl/WEB-INF/xsl/basic.xsl
%%APP_NAME%%/doc/examples/xsl/WEB-INF/xsl/content.xsl
%%APP_NAME%%/doc/examples/xsl/WEB-INF/xsl/copy.xsl
%%APP_NAME%%/doc/examples/xsl/WEB-INF/xsl/counter.xsl
%%APP_NAME%%/doc/examples/xsl/WEB-INF/xsl/default.xsl
%%APP_NAME%%/doc/examples/xsl/WEB-INF/xsl/header.xsl
%%APP_NAME%%/doc/examples/xsl/WEB-INF/xsl/plain.xsl
%%APP_NAME%%/doc/examples/xsl/WEB-INF/xsl/verify.xsl
%%APP_NAME%%/doc/examples/xsl/basic.xtp
%%APP_NAME%%/doc/examples/xsl/counter.xtp
%%APP_NAME%%/doc/examples/xsl/data/fortune/1.html
%%APP_NAME%%/doc/examples/xsl/data/fortune/2.html
%%APP_NAME%%/doc/examples/xsl/data/fortune/3.html
%%APP_NAME%%/doc/examples/xsl/data/fortune/4.html
%%APP_NAME%%/doc/examples/xsl/data/fortune/5.html
%%APP_NAME%%/doc/examples/xsl/good-doc.xtp
%%APP_NAME%%/doc/examples/xsl/index.xml
%%APP_NAME%%/doc/examples/xsl/navigation.xtp
%%APP_NAME%%/doc/examples/xsl/toc.xml
%%APP_NAME%%/doc/examples/xsl/vary.xml
%%APP_NAME%%/doc/examples/xsl/vary.xtp
%%APP_NAME%%/doc/examples/xsl/verify-bad.xtp
%%APP_NAME%%/doc/examples/xsl/verify-ok.xtp
%%APP_NAME%%/doc/faq/classloader.xtp
%%APP_NAME%%/doc/faq/config.xtp
%%APP_NAME%%/doc/faq/debug.xtp
%%APP_NAME%%/doc/faq/index.xtp
%%APP_NAME%%/doc/faq/install.xtp
%%APP_NAME%%/doc/faq/misc.xtp
%%APP_NAME%%/doc/faq/servlet.xtp
%%APP_NAME%%/doc/faq/ssl.xtp
%%APP_NAME%%/doc/faq/toc.xml
%%APP_NAME%%/doc/features.xtp
%%APP_NAME%%/doc/future.xtp
%%APP_NAME%%/doc/images/apache_srun.gif
%%APP_NAME%%/doc/images/appserver.gif
%%APP_NAME%%/doc/images/background.gif
%%APP_NAME%%/doc/images/backup.gif
%%APP_NAME%%/doc/images/ball8.gif
%%APP_NAME%%/doc/images/bean.gif
%%APP_NAME%%/doc/images/bean1.gif
%%APP_NAME%%/doc/images/caucho-blue.gif
%%APP_NAME%%/doc/images/caucho.gif
%%APP_NAME%%/doc/images/compiler.gif
%%APP_NAME%%/doc/images/file.gif
%%APP_NAME%%/doc/images/folder.gif
%%APP_NAME%%/doc/images/hbleed.gif
%%APP_NAME%%/doc/images/httpd.gif
%%APP_NAME%%/doc/images/iis_srun.gif
%%APP_NAME%%/doc/images/jndi_tree.gif
%%APP_NAME%%/doc/images/jsp.gif
%%APP_NAME%%/doc/images/left_background.gif
%%APP_NAME%%/doc/images/logo.gif
%%APP_NAME%%/doc/images/pixel.gif
%%APP_NAME%%/doc/images/resin_powered.gif
%%APP_NAME%%/doc/images/same_srun.gif
%%APP_NAME%%/doc/images/script.gif
%%APP_NAME%%/doc/images/srunc.gif
%%APP_NAME%%/doc/images/sticky-sessions.gif
%%APP_NAME%%/doc/images/style.gif
%%APP_NAME%%/doc/images/webserver.gif
%%APP_NAME%%/doc/images/website_srun.gif
%%APP_NAME%%/doc/index.xtp
%%APP_NAME%%/doc/java_tut/WEB-INF/classes/example/servlet/basic/HelloServlet.java
%%APP_NAME%%/doc/java_tut/WEB-INF/classes/example/servlet/basic/InitServlet.java
%%APP_NAME%%/doc/java_tut/WEB-INF/classes/example/servlet/basic/MappingServlet.java
%%APP_NAME%%/doc/java_tut/WEB-INF/classes/example/servlet/basic/URLServlet.java
%%APP_NAME%%/doc/java_tut/WEB-INF/classes/example/servlet/database/InitDatabaseServlet.java
%%APP_NAME%%/doc/java_tut/WEB-INF/classes/example/servlet/database/QueryServlet.java
%%APP_NAME%%/doc/java_tut/WEB-INF/sql/default.sql
%%APP_NAME%%/doc/java_tut/WEB-INF/web.xml
%%APP_NAME%%/doc/java_tut/add_page.xtp
%%APP_NAME%%/doc/java_tut/bean.xtp
%%APP_NAME%%/doc/java_tut/cache.xtp
%%APP_NAME%%/doc/java_tut/conf-basic.xtp
%%APP_NAME%%/doc/java_tut/conf-dir.xtp
%%APP_NAME%%/doc/java_tut/conf-servlet.xtp
%%APP_NAME%%/doc/java_tut/conf.xtp
%%APP_NAME%%/doc/java_tut/db-config.xtp
%%APP_NAME%%/doc/java_tut/db_form.xtp
%%APP_NAME%%/doc/java_tut/debug.xtp
%%APP_NAME%%/doc/java_tut/filter-vary.xtp
%%APP_NAME%%/doc/java_tut/form.xtp
%%APP_NAME%%/doc/java_tut/hardcore.xtp
%%APP_NAME%%/doc/java_tut/images/same_srun.gif
%%APP_NAME%%/doc/java_tut/images/srunc.gif
%%APP_NAME%%/doc/java_tut/index.xtp
%%APP_NAME%%/doc/java_tut/isp.xtp
%%APP_NAME%%/doc/java_tut/javadoc/allclasses-frame.html
%%APP_NAME%%/doc/java_tut/javadoc/allclasses-noframe.html
%%APP_NAME%%/doc/java_tut/javadoc/constant-values.html
%%APP_NAME%%/doc/java_tut/javadoc/deprecated-list.html
%%APP_NAME%%/doc/java_tut/javadoc/example/servlet/basic/HelloServlet.html
%%APP_NAME%%/doc/java_tut/javadoc/example/servlet/basic/InitServlet.html
%%APP_NAME%%/doc/java_tut/javadoc/example/servlet/basic/MappingServlet.html
%%APP_NAME%%/doc/java_tut/javadoc/example/servlet/basic/URLServlet.html
%%APP_NAME%%/doc/java_tut/javadoc/example/servlet/basic/package-frame.html
%%APP_NAME%%/doc/java_tut/javadoc/example/servlet/basic/package-summary.html
%%APP_NAME%%/doc/java_tut/javadoc/example/servlet/basic/package-tree.html
%%APP_NAME%%/doc/java_tut/javadoc/example/servlet/database/InitDatabaseServlet.html
%%APP_NAME%%/doc/java_tut/javadoc/example/servlet/database/QueryServlet.html
%%APP_NAME%%/doc/java_tut/javadoc/example/servlet/database/package-frame.html
%%APP_NAME%%/doc/java_tut/javadoc/example/servlet/database/package-summary.html
%%APP_NAME%%/doc/java_tut/javadoc/example/servlet/database/package-tree.html
%%APP_NAME%%/doc/java_tut/javadoc/help-doc.html
%%APP_NAME%%/doc/java_tut/javadoc/index-all.html
%%APP_NAME%%/doc/java_tut/javadoc/index.html
%%APP_NAME%%/doc/java_tut/javadoc/overview-frame.html
%%APP_NAME%%/doc/java_tut/javadoc/overview-summary.html
%%APP_NAME%%/doc/java_tut/javadoc/overview-tree.html
%%APP_NAME%%/doc/java_tut/javadoc/package-list
%%APP_NAME%%/doc/java_tut/javadoc/resources/inherit.gif
%%APP_NAME%%/doc/java_tut/javadoc/stylesheet.css
%%APP_NAME%%/doc/java_tut/jndi-appconfig.xtp
%%APP_NAME%%/doc/java_tut/linux-boot.xtp
%%APP_NAME%%/doc/java_tut/linux-init.xtp
%%APP_NAME%%/doc/java_tut/mail_form.xtp
%%APP_NAME%%/doc/java_tut/multipart.xtp
%%APP_NAME%%/doc/java_tut/parse_error.xtp
%%APP_NAME%%/doc/java_tut/plugin-dispatch.xtp
%%APP_NAME%%/doc/java_tut/req_headers.xtp
%%APP_NAME%%/doc/java_tut/req_url.xtp
%%APP_NAME%%/doc/java_tut/run-at.xtp
%%APP_NAME%%/doc/java_tut/run_error.xtp
%%APP_NAME%%/doc/java_tut/servlets.xtp
%%APP_NAME%%/doc/java_tut/tag-attribute.xtp
%%APP_NAME%%/doc/java_tut/tag-connection.xtp
%%APP_NAME%%/doc/java_tut/tag-hello.xtp
%%APP_NAME%%/doc/java_tut/tag-mail.xtp
%%APP_NAME%%/doc/java_tut/tag-xsl.xtp
%%APP_NAME%%/doc/java_tut/tcp-sessions.xtp
%%APP_NAME%%/doc/java_tut/toc.xml
%%APP_NAME%%/doc/java_tut/tuning.xtp
%%APP_NAME%%/doc/java_tut/war-dir.xtp
%%APP_NAME%%/doc/java_tut/xsl-filter.xtp
%%APP_NAME%%/doc/java_tut/xtp-copy.xtp
%%APP_NAME%%/doc/java_tut/xtp-format.xtp
%%APP_NAME%%/doc/java_tut/xtp-jsp-stylescript.xtp
%%APP_NAME%%/doc/java_tut/xtp-jsp.xtp
%%APP_NAME%%/doc/java_tut/xtp-page.xtp
%%APP_NAME%%/doc/java_tut/xtp-strict.xtp
%%APP_NAME%%/doc/ref/admin.xtp
%%APP_NAME%%/doc/ref/app-config.xtp
%%APP_NAME%%/doc/ref/atp.xtp
%%APP_NAME%%/doc/ref/auth-config.xtp
%%APP_NAME%%/doc/ref/balance.xtp
%%APP_NAME%%/doc/ref/basic-config.xtp
%%APP_NAME%%/doc/ref/bnf.xtp
%%APP_NAME%%/doc/ref/burlap.xtp
%%APP_NAME%%/doc/ref/caching.xtp
%%APP_NAME%%/doc/ref/changes.xtp
%%APP_NAME%%/doc/ref/config-index.xtp
%%APP_NAME%%/doc/ref/config-sum.xtp
%%APP_NAME%%/doc/ref/config.xtp
%%APP_NAME%%/doc/ref/contents.xtp
%%APP_NAME%%/doc/ref/cse-apache-unix.xtp
%%APP_NAME%%/doc/ref/cse-apache-win.xtp
%%APP_NAME%%/doc/ref/cse-iis.xtp
%%APP_NAME%%/doc/ref/cse-netscape.xtp
%%APP_NAME%%/doc/ref/cse-oreilly-website.xtp
%%APP_NAME%%/doc/ref/db-config.xtp
%%APP_NAME%%/doc/ref/ejb-ias.xtp
%%APP_NAME%%/doc/ref/ejb-jboss.xtp
%%APP_NAME%%/doc/ref/ejb-jonas.xtp
%%APP_NAME%%/doc/ref/ejb-orion.xtp
%%APP_NAME%%/doc/ref/ejb-weblogic.xtp
%%APP_NAME%%/doc/ref/ejb.xtp
%%APP_NAME%%/doc/ref/errata.xtp
%%APP_NAME%%/doc/ref/faq.xtp
%%APP_NAME%%/doc/ref/file-config.xtp
%%APP_NAME%%/doc/ref/filesystems.xtp
%%APP_NAME%%/doc/ref/filter-config.xtp
%%APP_NAME%%/doc/ref/filter-library.xtp
%%APP_NAME%%/doc/ref/filter.xtp
%%APP_NAME%%/doc/ref/glossary.xtp
%%APP_NAME%%/doc/ref/hessian.xtp
%%APP_NAME%%/doc/ref/http-config.xtp
%%APP_NAME%%/doc/ref/httpd.xtp
%%APP_NAME%%/doc/ref/ide-jbuilder.xtp
%%APP_NAME%%/doc/ref/ide-kawa.xtp
%%APP_NAME%%/doc/ref/ide.xtp
%%APP_NAME%%/doc/ref/index.xtp
%%APP_NAME%%/doc/ref/javascript.xtp
%%APP_NAME%%/doc/ref/jms.xtp
%%APP_NAME%%/doc/ref/jmx.xtp
%%APP_NAME%%/doc/ref/jndi.xtp
%%APP_NAME%%/doc/ref/jsp/jsp-el.xtp
%%APP_NAME%%/doc/ref/jsp-actions.xtp
%%APP_NAME%%/doc/ref/jsp-applications.xtp
%%APP_NAME%%/doc/ref/jsp-directives.xtp
%%APP_NAME%%/doc/ref/jsp-variables.xtp
%%APP_NAME%%/doc/ref/jsp.xtp
%%APP_NAME%%/doc/ref/jspapp.xtp
%%APP_NAME%%/doc/ref/jtp.xtp
%%APP_NAME%%/doc/ref/log-config.xtp
%%APP_NAME%%/doc/ref/macos-x.xtp
%%APP_NAME%%/doc/ref/page_context.xtp
%%APP_NAME%%/doc/ref/port-config.xtp
%%APP_NAME%%/doc/ref/raw-socket.xtp
%%APP_NAME%%/doc/ref/resin-config.xtp
%%APP_NAME%%/doc/ref/security.xtp
%%APP_NAME%%/doc/ref/servlet-config.xtp
%%APP_NAME%%/doc/ref/servlet-library.xtp
%%APP_NAME%%/doc/ref/servlet.xtp
%%APP_NAME%%/doc/ref/servletcontext.xtp
%%APP_NAME%%/doc/ref/sessions.xtp
%%APP_NAME%%/doc/ref/srun-other.xtp
%%APP_NAME%%/doc/ref/ssl-quick.xml
%%APP_NAME%%/doc/ref/ssl-quick.xtp
%%APP_NAME%%/doc/ref/ssl.xml
%%APP_NAME%%/doc/ref/ssl.xtp
%%APP_NAME%%/doc/ref/starting.xtp
%%APP_NAME%%/doc/ref/static.xtp
%%APP_NAME%%/doc/ref/stylescript.xtp
%%APP_NAME%%/doc/ref/taglib.xtp
%%APP_NAME%%/doc/ref/toc.xml
%%APP_NAME%%/doc/ref/topics.xtp
%%APP_NAME%%/doc/ref/trouble.xtp
%%APP_NAME%%/doc/ref/velocity.xtp
%%APP_NAME%%/doc/ref/virtual-host.xtp
%%APP_NAME%%/doc/ref/visual-age.xtp
%%APP_NAME%%/doc/ref/war-dir.xtp
%%APP_NAME%%/doc/ref/webdav.xtp
%%APP_NAME%%/doc/ref/xml/attributes.xtp
%%APP_NAME%%/doc/ref/xml/elements.xtp
%%APP_NAME%%/doc/ref/xml/get_elements.xtp
%%APP_NAME%%/doc/ref/xml/html.xtp
%%APP_NAME%%/doc/ref/xml/index.xtp
%%APP_NAME%%/doc/ref/xml/intro.xtp
%%APP_NAME%%/doc/ref/xml/parsing.xtp
%%APP_NAME%%/doc/ref/xml/sax.xtp
%%APP_NAME%%/doc/ref/xml/text.xtp
%%APP_NAME%%/doc/ref/xml/vfs.xtp
%%APP_NAME%%/doc/ref/xml.xtp
%%APP_NAME%%/doc/ref/xpath-fun.xtp
%%APP_NAME%%/doc/ref/xpath.xtp
%%APP_NAME%%/doc/ref/xsl/apply-templates.xtp
%%APP_NAME%%/doc/ref/xsl/attribute.xtp
%%APP_NAME%%/doc/ref/xsl/boxes.xtp
%%APP_NAME%%/doc/ref/xsl/command-line.xtp
%%APP_NAME%%/doc/ref/xsl/copy.xtp
%%APP_NAME%%/doc/ref/xsl/filter-attribute.xtp
%%APP_NAME%%/doc/ref/xsl/filter.xtp
%%APP_NAME%%/doc/ref/xsl/hello.xtp
%%APP_NAME%%/doc/ref/xsl/if-attribute.xtp
%%APP_NAME%%/doc/ref/xsl/if-xpath-attribute.xtp
%%APP_NAME%%/doc/ref/xsl/if.xtp
%%APP_NAME%%/doc/ref/xsl/index.xtp
%%APP_NAME%%/doc/ref/xsl/modes.xtp
%%APP_NAME%%/doc/ref/xsl/ref.xtp
%%APP_NAME%%/doc/ref/xsl/short-attribute.xtp
%%APP_NAME%%/doc/ref/xsl/value-of.xtp
%%APP_NAME%%/doc/ref/xsl/xpath-find.xtp
%%APP_NAME%%/doc/ref/xsl/xpath-select.xtp
%%APP_NAME%%/doc/ref/xsl/xpath.xtp
%%APP_NAME%%/doc/ref/xsl/xsl-api.xtp
%%APP_NAME%%/doc/ref/xsl/xsl-boxes.xtp
%%APP_NAME%%/doc/ref/xsl/xsl.xtp
%%APP_NAME%%/doc/ref/xsl.xtp
%%APP_NAME%%/doc/ref/xslt-lite-eg.xtp
%%APP_NAME%%/doc/ref/xslt-lite.xtp
%%APP_NAME%%/doc/ref/xtp-index.xtp
%%APP_NAME%%/doc/ref/xtp-intro.xtp
%%APP_NAME%%/doc/ref/xtp-jsp.xtp
%%APP_NAME%%/doc/ref/xtp.xtp
%%APP_NAME%%/doc/toc.xml
%%APP_NAME%%/lib/dom.jar
%%APP_NAME%%/lib/jaxp.jar
%%APP_NAME%%/lib/jdbc-mysql.jar
%%APP_NAME%%/lib/jdbc2_0-stdext.jar
%%APP_NAME%%/lib/jmx.jar
%%APP_NAME%%/lib/jndi.jar
%%APP_NAME%%/lib/jsdk23.jar
%%APP_NAME%%/lib/jta_101.jar
%%APP_NAME%%/lib/resin.jar
%%APP_NAME%%/lib/sax.jar
%%APP_NAME%%/lib/webutil.jar
%%APP_NAME%%/xsl/caucho/comment.xsl
%%APP_NAME%%/xsl/caucho/contents.xsl
%%APP_NAME%%/xsl/caucho/control.xsl
%%APP_NAME%%/xsl/caucho/debug.xsl
%%APP_NAME%%/xsl/caucho/defaultcopy.xsl
%%APP_NAME%%/xsl/caucho/format.xsl
%%APP_NAME%%/xsl/caucho/newstyle.xsl
%%APP_NAME%%/xsl/caucho/weblog.xsl
%%APP_NAME%%/xsl/caucho_default.xsl
%%APP_NAME%%/xsl/caucho_header.xsl
%%APP_NAME%%/xsl/caucho_subpage.xsl
%%APP_NAME%%/xsl/default.xsl
%%APP_NAME%%/xsl/doc_header.xsl
%%APP_NAME%%/xsl/header.xsl
%%APP_NAME%%/xsl/js/caucho/defaultcopy.xsl
%%APP_NAME%%/xsl/js/caucho/format.xsl
%%APP_NAME%%/xsl/js/caucho/weblog.xsl
%%APP_NAME%%/xsl/js_header.xsl
%%APP_NAME%%/xsl/new_doc_header.xsl
%%APP_NAME%%/xsl/search.xsl
%%APP_NAME%%/xsl/subpage.xsl
%%APP_NAME%%/libexec/libresin.so
etc/rc.d/%%APP_NAME%%.sh
sbin/%%APP_NAME%%ctl
%%MOD_DIR%%/mod_caucho.so
%%APACHE%%@exec %%APXS%% -e -a -n caucho %f
%%APACHE%%@unexec %%APXS%% -e -A -n caucho %f
@dirrm %%APP_NAME%%/doc/css
@dirrm %%APP_NAME%%/doc/examples/basic/WEB-INF/classes
@dirrm %%APP_NAME%%/doc/examples/basic/WEB-INF
@dirrm %%APP_NAME%%/doc/examples/basic/com/caucho/fortune
@dirrm %%APP_NAME%%/doc/examples/basic/com/caucho
@dirrm %%APP_NAME%%/doc/examples/basic/com
@dirrm %%APP_NAME%%/doc/examples/basic/data/fortune
@dirrm %%APP_NAME%%/doc/examples/basic/data
@dirrm %%APP_NAME%%/doc/examples/basic
@dirrm %%APP_NAME%%/doc/examples/login/WEB-INF/classes/test
@dirrm %%APP_NAME%%/doc/examples/login/WEB-INF/classes
@dirrm %%APP_NAME%%/doc/examples/login/WEB-INF
@dirrm %%APP_NAME%%/doc/examples/login
@dirrm %%APP_NAME%%/doc/examples/navigation/WEB-INF/scripts
@dirrm %%APP_NAME%%/doc/examples/navigation/WEB-INF
@dirrm %%APP_NAME%%/doc/examples/navigation
@dirrm %%APP_NAME%%/doc/examples/templates/WEB-INF/classes/jsp
@dirrm %%APP_NAME%%/doc/examples/templates/WEB-INF/classes
@dirrm %%APP_NAME%%/doc/examples/templates/WEB-INF/xsl
@dirrm %%APP_NAME%%/doc/examples/templates/WEB-INF
@dirrm %%APP_NAME%%/doc/examples/templates
@dirrm %%APP_NAME%%/doc/examples/tictactoe/WEB-INF/scripts
@dirrm %%APP_NAME%%/doc/examples/tictactoe/WEB-INF
@dirrm %%APP_NAME%%/doc/examples/tictactoe
@dirrm %%APP_NAME%%/doc/examples/tutorial/WEB-INF/classes/servlet/trivial
@dirrm %%APP_NAME%%/doc/examples/tutorial/WEB-INF/classes/servlet
@dirrm %%APP_NAME%%/doc/examples/tutorial/WEB-INF/classes
@dirrm %%APP_NAME%%/doc/examples/tutorial/WEB-INF
@dirrm %%APP_NAME%%/doc/examples/tutorial
@dirrm %%APP_NAME%%/doc/examples/xsl/WEB-INF/classes/xsl
@dirrm %%APP_NAME%%/doc/examples/xsl/WEB-INF/classes
@dirrm %%APP_NAME%%/doc/examples/xsl/WEB-INF/scripts
@dirrm %%APP_NAME%%/doc/examples/xsl/WEB-INF/xsl
@dirrm %%APP_NAME%%/doc/examples/xsl/WEB-INF
@dirrm %%APP_NAME%%/doc/examples/xsl/data/fortune
@dirrm %%APP_NAME%%/doc/examples/xsl/data
@dirrm %%APP_NAME%%/doc/examples/xsl
@dirrm %%APP_NAME%%/doc/examples
@dirrm %%APP_NAME%%/doc/faq
@dirrm %%APP_NAME%%/doc/images
@dirrm %%APP_NAME%%/doc/java_tut/WEB-INF/classes/example/servlet/basic
@dirrm %%APP_NAME%%/doc/java_tut/WEB-INF/classes/example/servlet/database
@dirrm %%APP_NAME%%/doc/java_tut/WEB-INF/classes/example/servlet
@dirrm %%APP_NAME%%/doc/java_tut/WEB-INF/classes/example
@dirrm %%APP_NAME%%/doc/java_tut/WEB-INF/classes
@dirrm %%APP_NAME%%/doc/java_tut/WEB-INF/sql
@dirrm %%APP_NAME%%/doc/java_tut/WEB-INF
@dirrm %%APP_NAME%%/doc/java_tut/images
@dirrm %%APP_NAME%%/doc/java_tut/javadoc/example/servlet/basic
@dirrm %%APP_NAME%%/doc/java_tut/javadoc/example/servlet/database
@dirrm %%APP_NAME%%/doc/java_tut/javadoc/example/servlet
@dirrm %%APP_NAME%%/doc/java_tut/javadoc/example
@dirrm %%APP_NAME%%/doc/java_tut/javadoc/resources
@dirrm %%APP_NAME%%/doc/java_tut/javadoc
@dirrm %%APP_NAME%%/doc/java_tut
@dirrm %%APP_NAME%%/doc/ref/jsp
@dirrm %%APP_NAME%%/doc/ref/xml
@dirrm %%APP_NAME%%/doc/ref/xsl
@dirrm %%APP_NAME%%/doc/ref
@dirrm %%APP_NAME%%/doc
@dirrm %%APP_NAME%%/lib
@dirrm %%APP_NAME%%/xsl/caucho
@dirrm %%APP_NAME%%/xsl/js/caucho
@dirrm %%APP_NAME%%/xsl/js
@dirrm %%APP_NAME%%/xsl
@dirrm %%APP_NAME%%/libexec
@dirrm %%APP_NAME%%
@unexec if cmp -s %D/etc/%%APP_NAME%%.xml %D/etc/%%APP_NAME%%.xml-dist; then rm -f %D/etc/%%APP_NAME%%.xml; else echo If permanently deleting this package, %D/etc/%%APP_NAME%%.xml must be removed manually; fi
etc/%%APP_NAME%%.xml-dist
@exec if test -f %B/%%APP_NAME%%.xml ; then echo Keeping %B/%%APP_NAME%%.xml intact from previous installation; else cp %F %B/%%APP_NAME%%.xml; fi
