# misc/e2fsprogs-libuuid/pkg-plist
bin/uuidgen
include/uuid/uuid.h
lib/libuuid.so.1.2
lib/libuuid.so.1
lib/libuuid.so
lib/libuuid.a
libdata/pkgconfig/uuid.pc
man/man1/uuidgen.1.gz
man/man3/uuid.3.gz
man/man3/uuid_clear.3.gz
man/man3/uuid_compare.3.gz
man/man3/uuid_copy.3.gz
man/man3/uuid_generate.3.gz
man/man3/uuid_is_null.3.gz
man/man3/uuid_parse.3.gz
man/man3/uuid_time.3.gz
man/man3/uuid_unparse.3.gz
man/man3/uuid_generate_random.3.gz
man/man3/uuid_generate_time.3.gz
man/man8/uuidd.8.gz
sbin/uuidd
@dir /var/run/libuuid
@preunexec %D/etc/rc.d/uuidd stop  >/dev/null 2>&1 || true
