# sysutils/moosefs-chunkserver/files/mfschunkserver.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: mfschunkserver
# REQUIRE: LOGIN 
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf to enable mfschunkserver:
#
# mfschunkserver_enable (bool):   Set to NO by default.
#                                 Set it to YES to enable mfschunkserver.
# mfschunkserver_config (path):   Set to %%PREFIX%%/etc/mfs/mfschunkserver.cfg
#                                 by default.
# mfschunkserver_flags (str):     Set to ""  by default.
#                                 Extra flags passed to mfschunkserver.
#

. /etc/rc.subr

name=mfschunkserver
rcvar=mfschunkserver_enable

load_rc_config $name

: ${mfschunkserver_enable=NO}
: ${mfschunkserver_config=%%PREFIX%%/etc/mfs/mfschunkserver.cfg}

extra_commands=reload
command=%%PREFIX%%/sbin/${name}
command_args="-c ${mfschunkserver_config}"

run_rc_command "$1"
