# emulators/gngeo/pkg-plist
bin/gngeo
%%EXAMPLESDIR%%/gngeorc
%%DATADIR%%/cursor.bmp
%%DATADIR%%/gui_font.bmp
%%DATADIR%%/gui_font2.bmp
%%DATADIR%%/little_font.bmp
%%DATADIR%%/romrc
%%DATADIR%%/romrc.d/2020bb.rc
%%DATADIR%%/romrc.d/2020bba.rc
%%DATADIR%%/romrc.d/2020bbh.rc
%%DATADIR%%/romrc.d/3countb.rc
%%DATADIR%%/romrc.d/alpham2.rc
%%DATADIR%%/romrc.d/androdun.rc
%%DATADIR%%/romrc.d/aodk.rc
%%DATADIR%%/romrc.d/aof.rc
%%DATADIR%%/romrc.d/aof2.rc
%%DATADIR%%/romrc.d/aof2a.rc
%%DATADIR%%/romrc.d/aof3.rc
%%DATADIR%%/romrc.d/bakatono.rc
%%DATADIR%%/romrc.d/bangbead.rc
%%DATADIR%%/romrc.d/bjourney.rc
%%DATADIR%%/romrc.d/blazstar.rc
%%DATADIR%%/romrc.d/breakers.rc
%%DATADIR%%/romrc.d/breakrev.rc
%%DATADIR%%/romrc.d/bstars.rc
%%DATADIR%%/romrc.d/bstars2.rc
%%DATADIR%%/romrc.d/burningf.rc
%%DATADIR%%/romrc.d/burningh.rc
%%DATADIR%%/romrc.d/crsword.rc
%%DATADIR%%/romrc.d/cthd2003.rc
%%DATADIR%%/romrc.d/cthd2003sp.rc
%%DATADIR%%/romrc.d/ctomaday.rc
%%DATADIR%%/romrc.d/cyberlip.rc
%%DATADIR%%/romrc.d/doubledr.rc
%%DATADIR%%/romrc.d/eightman.rc
%%DATADIR%%/romrc.d/fatfursa.rc
%%DATADIR%%/romrc.d/fatfursp.rc
%%DATADIR%%/romrc.d/fatfury1.rc
%%DATADIR%%/romrc.d/fatfury2.rc
%%DATADIR%%/romrc.d/fatfury3.rc
%%DATADIR%%/romrc.d/fbfrenzy.rc
%%DATADIR%%/romrc.d/fightfev.rc
%%DATADIR%%/romrc.d/fightfva.rc
%%DATADIR%%/romrc.d/flipshot.rc
%%DATADIR%%/romrc.d/galaxyfg.rc
%%DATADIR%%/romrc.d/garou.rc
%%DATADIR%%/romrc.d/garoup.rc
%%DATADIR%%/romrc.d/ghostlop.rc
%%DATADIR%%/romrc.d/goalx3.rc
%%DATADIR%%/romrc.d/gowcaizr.rc
%%DATADIR%%/romrc.d/gpilots.rc
%%DATADIR%%/romrc.d/gururin.rc
%%DATADIR%%/romrc.d/janshin.rc
%%DATADIR%%/romrc.d/jockeygp.rc
%%DATADIR%%/romrc.d/joyjoy.rc
%%DATADIR%%/romrc.d/kabukikl.rc
%%DATADIR%%/romrc.d/karnovr.rc
%%DATADIR%%/romrc.d/kizuna.rc
%%DATADIR%%/romrc.d/kof10thu.rc
%%DATADIR%%/romrc.d/kof2000n.rc
%%DATADIR%%/romrc.d/kof2001n.rc
%%DATADIR%%/romrc.d/kof2002.rc
%%DATADIR%%/romrc.d/kof2003.rc
%%DATADIR%%/romrc.d/kof2003a.rc
%%DATADIR%%/romrc.d/kof2003b.rc
%%DATADIR%%/romrc.d/kof2k1rp.rc
%%DATADIR%%/romrc.d/kof2k2.rc
%%DATADIR%%/romrc.d/kof2k2pls.rc
%%DATADIR%%/romrc.d/kof2k4es.rc
%%DATADIR%%/romrc.d/kof94.rc
%%DATADIR%%/romrc.d/kof95.rc
%%DATADIR%%/romrc.d/kof95a.rc
%%DATADIR%%/romrc.d/kof96.rc
%%DATADIR%%/romrc.d/kof96h.rc
%%DATADIR%%/romrc.d/kof97.rc
%%DATADIR%%/romrc.d/kof97a.rc
%%DATADIR%%/romrc.d/kof97pls.rc
%%DATADIR%%/romrc.d/kof98.rc
%%DATADIR%%/romrc.d/kof98k.rc
%%DATADIR%%/romrc.d/kof98n.rc
%%DATADIR%%/romrc.d/kof99.rc
%%DATADIR%%/romrc.d/kof99p.rc
%%DATADIR%%/romrc.d/kog.rc
%%DATADIR%%/romrc.d/kotm.rc
%%DATADIR%%/romrc.d/kotm2.rc
%%DATADIR%%/romrc.d/kotmh.rc
%%DATADIR%%/romrc.d/lastblad.rc
%%DATADIR%%/romrc.d/lastbld2.rc
%%DATADIR%%/romrc.d/lastblda.rc
%%DATADIR%%/romrc.d/lbowling.rc
%%DATADIR%%/romrc.d/legendos.rc
%%DATADIR%%/romrc.d/lresort.rc
%%DATADIR%%/romrc.d/magdrop2.rc
%%DATADIR%%/romrc.d/magdrop3.rc
%%DATADIR%%/romrc.d/maglord.rc
%%DATADIR%%/romrc.d/maglordh.rc
%%DATADIR%%/romrc.d/mahretsu.rc
%%DATADIR%%/romrc.d/marukodq.rc
%%DATADIR%%/romrc.d/miexchng.rc
%%DATADIR%%/romrc.d/minasan.rc
%%DATADIR%%/romrc.d/mosyougi.rc
%%DATADIR%%/romrc.d/mslug.rc
%%DATADIR%%/romrc.d/mslug2.rc
%%DATADIR%%/romrc.d/mslug3.rc
%%DATADIR%%/romrc.d/mslug3n.rc
%%DATADIR%%/romrc.d/mslug4.rc
%%DATADIR%%/romrc.d/mslug5.rc
%%DATADIR%%/romrc.d/mslug5nd.rc
%%DATADIR%%/romrc.d/mslugx.rc
%%DATADIR%%/romrc.d/mutnat.rc
%%DATADIR%%/romrc.d/nam1975.rc
%%DATADIR%%/romrc.d/ncombat.rc
%%DATADIR%%/romrc.d/ncombata.rc
%%DATADIR%%/romrc.d/ncommand.rc
%%DATADIR%%/romrc.d/neobombe.rc
%%DATADIR%%/romrc.d/neocup98.rc
%%DATADIR%%/romrc.d/neodrift.rc
%%DATADIR%%/romrc.d/neomrdo.rc
%%DATADIR%%/romrc.d/neopong.rc
%%DATADIR%%/romrc.d/ninjamas.rc
%%DATADIR%%/romrc.d/nitd.rc
%%DATADIR%%/romrc.d/overtop.rc
%%DATADIR%%/romrc.d/panicbom.rc
%%DATADIR%%/romrc.d/pbobbl2n.rc
%%DATADIR%%/romrc.d/pbobble.rc
%%DATADIR%%/romrc.d/pbobblen.rc
%%DATADIR%%/romrc.d/pbobblna.rc
%%DATADIR%%/romrc.d/pgoal.rc
%%DATADIR%%/romrc.d/pim.rc
%%DATADIR%%/romrc.d/pochi.rc
%%DATADIR%%/romrc.d/popbounc.rc
%%DATADIR%%/romrc.d/preisle2.rc
%%DATADIR%%/romrc.d/pspikes2.rc
%%DATADIR%%/romrc.d/pulstar.rc
%%DATADIR%%/romrc.d/puzzldpr.rc
%%DATADIR%%/romrc.d/puzzledp.rc
%%DATADIR%%/romrc.d/quizdai2.rc
%%DATADIR%%/romrc.d/quizdais.rc
%%DATADIR%%/romrc.d/quizkof.rc
%%DATADIR%%/romrc.d/ragnagrd.rc
%%DATADIR%%/romrc.d/rbff1.rc
%%DATADIR%%/romrc.d/rbff2.rc
%%DATADIR%%/romrc.d/rbff2a.rc
%%DATADIR%%/romrc.d/rbffspec.rc
%%DATADIR%%/romrc.d/ridhero.rc
%%DATADIR%%/romrc.d/ridheroh.rc
%%DATADIR%%/romrc.d/roboarma.rc
%%DATADIR%%/romrc.d/roboarmy.rc
%%DATADIR%%/romrc.d/rodd.rc
%%DATADIR%%/romrc.d/s1945p.rc
%%DATADIR%%/romrc.d/samsho.rc
%%DATADIR%%/romrc.d/samsho2.rc
%%DATADIR%%/romrc.d/samsho3.rc
%%DATADIR%%/romrc.d/samsho3a.rc
%%DATADIR%%/romrc.d/samsho4.rc
%%DATADIR%%/romrc.d/samsho5.rc
%%DATADIR%%/romrc.d/samsho5sp.rc
%%DATADIR%%/romrc.d/savagere.rc
%%DATADIR%%/romrc.d/sdodgeb.rc
%%DATADIR%%/romrc.d/sengokh.rc
%%DATADIR%%/romrc.d/sengoku.rc
%%DATADIR%%/romrc.d/sengoku2.rc
%%DATADIR%%/romrc.d/sengoku3.rc
%%DATADIR%%/romrc.d/shocktr2.rc
%%DATADIR%%/romrc.d/shocktra.rc
%%DATADIR%%/romrc.d/shocktrj.rc
%%DATADIR%%/romrc.d/shocktro.rc
%%DATADIR%%/romrc.d/socbrawl.rc
%%DATADIR%%/romrc.d/sonicwi2.rc
%%DATADIR%%/romrc.d/sonicwi3.rc
%%DATADIR%%/romrc.d/spinmast.rc
%%DATADIR%%/romrc.d/ssideki.rc
%%DATADIR%%/romrc.d/ssideki2.rc
%%DATADIR%%/romrc.d/ssideki3.rc
%%DATADIR%%/romrc.d/ssideki4.rc
%%DATADIR%%/romrc.d/ssvsp.rc
%%DATADIR%%/romrc.d/stakwin.rc
%%DATADIR%%/romrc.d/stakwin2.rc
%%DATADIR%%/romrc.d/strhoop.rc
%%DATADIR%%/romrc.d/superspy.rc
%%DATADIR%%/romrc.d/svcplus.rc
%%DATADIR%%/romrc.d/tophunta.rc
%%DATADIR%%/romrc.d/tophuntr.rc
%%DATADIR%%/romrc.d/tpgolf.rc
%%DATADIR%%/romrc.d/trally.rc
%%DATADIR%%/romrc.d/turfmast.rc
%%DATADIR%%/romrc.d/twinspri.rc
%%DATADIR%%/romrc.d/tws96.rc
%%DATADIR%%/romrc.d/viewpoin.rc
%%DATADIR%%/romrc.d/wakuwak7.rc
%%DATADIR%%/romrc.d/wh1.rc
%%DATADIR%%/romrc.d/wh1h.rc
%%DATADIR%%/romrc.d/wh2.rc
%%DATADIR%%/romrc.d/wh2j.rc
%%DATADIR%%/romrc.d/whp.rc
%%DATADIR%%/romrc.d/wjammers.rc
%%DATADIR%%/romrc.d/zedblade.rc
%%DATADIR%%/romrc.d/zintrckb.rc
%%DATADIR%%/romrc.d/zupapa.rc
%%DATADIR%%/sb_arrows.bmp
@dirrm %%EXAMPLESDIR%%
@dirrm %%DATADIR%%/romrc.d
@unexec rmdir %D/%%DATADIR%% >/dev/null 2>&1 || echo "If you are permanently removing this port, you should remove the Neo-Geo BIOS files left in %D/%%DATADIR%%."
