# www/thttpd/files/openrc-thttpd.in
#!/sbin/openrc-run

name=thttpd

command=%%PREFIX%%/sbin/thttpd

thttpd_flags=${thttpd_flags:-"-C %%PREFIX%%/etc/thttpd.conf"}

depend(){
	keyword -shutdown
}
