# www/caddy/files/caddy.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: caddy
# REQUIRE: LOGIN DAEMON NETWORKING
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
# caddy_enable (bool):   Set to NO by default. Set it to YES to enable caddy.
#
# caddy_config (string):      (Optional) Full path to caddy config file
# caddy_adapter (string):     (Optional) Adapter type if the configuration is not in caddyfile format
# caddy_extra_flags (string): (Optional) Flags passed to caddy start
# caddy_logfile (string):     Location of process log. This is for startup/shutdown/error messages.
#                             To create an access log, see https://caddyserver.com/docs/caddyfile/directives/log

. /etc/rc.subr

name=caddy
rcvar=caddy_enable
desc="Powerful, enterprise-ready, open source web server with automatic HTTPS written in Go"

load_rc_config $name

# Defaults
: ${caddy_enable:=NO}
: ${caddy_config:=%%ETCDIR%%/Caddyfile}
: ${caddy_adapter:=caddyfile}
: ${caddy_extra_flags:=""}
: ${caddy_logfile="/var/log/caddy.log"}

command="%%PREFIX%%/bin/${name}"
caddy_flags="--config ${caddy_config} --adapter ${caddy_adapter}"
pidfile="/var/run/${name}.pid"

required_files="${caddy_config} ${command}"

# Extra Commands
extra_commands="configtest reload"

configtest_cmd="${command} validate ${caddy_flags}"
reload_cmd="${command} reload ${caddy_flags}"
start_cmd="${command} start ${caddy_flags} ${caddy_extra_flags} --pidfile ${pidfile} >> ${caddy_logfile} 2>&1"
stop_cmd="${command} stop"

run_rc_command "$1"
