# www/grafana4/files/grafana4.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: grafana4
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable grafana
# grafana_enable="YES"
#
# grafana4_enable (bool): 	Set to YES to enable grafana
#				Default: NO
# grafana4_conf (str):		grafana configration file
#				Default: %%PREFIX%%/etc/grafana.cfg
# grafana4_user (str):		grafana daemon user
#				Default: %%GRAFANA_USER%%
# grafana4_group (str):		grafana daemon group
#				Default: %%GRAFANA_GROUP%%

. /etc/rc.subr

name="grafana4"
rcvar=grafana4_enable
load_rc_config $name

: ${grafana4_enable:="NO"}
: ${grafana4_user:="%%GRAFANA_USER%%"}
: ${grafana4_group:="%%GRAFANA_GROUP%%"}
: ${grafana4_flags:=""}
: ${grafana4_conf:="%%PREFIX%%/etc/${name}.conf"}
: ${grafana4_options:="-config=${grafana4_conf} -homepath=%%GRAFANAHOMEDIR%%"}

pidfile="%%GRAFANAPIDDIR%%${name}.pid"
procname="%%PREFIX%%/bin/grafana-server"
command=/usr/sbin/daemon
command_args="-f -c -p ${pidfile} ${procname} ${grafana4_options}"
start_precmd="grafana_precmd"

grafana_precmd()
{
    chown ${grafana4_user}:${grafana4_group} %%GRAFANADATADIR%%
    chown ${grafana4_user}:${grafana4_group} %%GRAFANALOGDIR%%
    chown ${grafana4_user}:${grafana4_group} %%GRAFANAHOMEDIR%%data/log/
    install -d -o ${grafana4_user} %%GRAFANAPIDDIR%%
    install -o ${grafana4_user} /dev/null ${pidfile}
}

run_rc_command "$1"
