# sysutils/rsyslog5-devel/files/rsyslogd.in
#!/bin/sh
#
# $FreeBSD: /tmp/pcvs/ports/sysutils/rsyslog5-devel/files/Attic/rsyslogd.in,v 1.4 2010-03-27 00:15:04 dougb Exp $
#


# PROVIDE: rsyslogd
# REQUIRE: mountcritremote cleanvar newsyslog ldconfig
# BEFORE:  SERVERS

. /etc/rc.subr

name=rsyslogd
rcvar=`set_rcvar`
command="%%PREFIX%%/sbin/${name}"
load_rc_config $name
: ${rsyslogd_enable:="NO"}
: ${rsyslogd_pidfile:="/var/run/rsyslogd.pid"}
: ${rsyslogd_config:="%%PREFIX%%/etc/rsyslog.conf"}
pidfile="${rsyslogd_pidfile}"
command_args="-i ${pidfile} -f ${rsyslogd_config}"
required_files="${rsyslogd_config}"
extra_commands="reload"

run_rc_command "$1"
