# news/inn/pkg-plist
etc/rc.d/innd.sh
news/bin/actmerge
news/bin/actsync
news/bin/actsyncd
news/bin/archive
news/bin/auth/passwd/ckpasswd
news/bin/auth/passwd/pwcheck
news/bin/auth/passwd/radius
news/bin/auth/resolv/domain
news/bin/auth/resolv/ident
news/bin/batcher
news/bin/buffchan
news/bin/cnfsheadconf
news/bin/cnfsstat
news/bin/control/checkgroups
news/bin/control/checkgroups.pl
news/bin/control/default
news/bin/control/ihave
news/bin/control/ihave.pl
news/bin/control/newgroup
news/bin/control/newgroup.pl
news/bin/control/rmgroup
news/bin/control/rmgroup.pl
news/bin/control/sendme
news/bin/control/sendme.pl
news/bin/control/sendsys
news/bin/control/sendsys.pl
news/bin/control/senduuname
news/bin/control/senduuname.pl
news/bin/control/version
news/bin/control/version.pl
news/bin/controlbatch
news/bin/controlchan
news/bin/convdate
news/bin/ctlinnd
news/bin/cvtbatch
news/bin/dbprocs
news/bin/expire
news/bin/expireover
news/bin/expirerm
news/bin/fastrm
news/bin/filechan
news/bin/filter/INN.py
news/bin/filter/filter.tcl
news/bin/filter/filter_innd.pl
news/bin/filter/filter_innd.py
news/bin/filter/filter_nnrpd.pl
news/bin/filter/nnrpd_auth.pl
news/bin/filter/nnrpd_auth.py
news/bin/filter/startup.tcl
news/bin/filter/startup_innd.pl
news/bin/getlist
news/bin/gpgverify
news/bin/grephistory
news/bin/inews
news/bin/inncheck
news/bin/innconfval
news/bin/innd
news/bin/inndf
news/bin/inndstart
news/bin/innfeed
news/bin/innmail
news/bin/innreport
news/bin/innstat
news/bin/innwatch
news/bin/innxbatch
news/bin/innxmit
news/bin/mailpost
news/bin/makedbz
news/bin/makehistory
news/bin/mod-active
news/bin/news.daily
news/bin/news2mail
news/bin/newsrequeue
news/bin/nnrpd
news/bin/nntpget
news/bin/nntpsend
news/bin/ovdb_recover
news/bin/ovdb_upgrade
news/bin/overchan
news/bin/parsecontrol
news/bin/pgpverify
news/bin/procbatch
news/bin/prunehistory
news/bin/pullnews
news/bin/rc.news
news/bin/rnews
news/bin/rnews.libexec/c7unbatch
news/bin/rnews.libexec/decode
news/bin/rnews.libexec/encode
news/bin/rnews.libexec/gunbatch
news/bin/scanlogs
news/bin/scanspool
news/bin/send-ihave
news/bin/send-nntp
news/bin/send-uucp
news/bin/sendbatch
news/bin/sendxbatches
news/bin/shlock
news/bin/shrinkfile
news/bin/signcontrol
news/bin/simpleftp
news/bin/sm
news/bin/startinnfeed
news/bin/tally.control
news/bin/writelog
news/db/active
news/db/active.times
news/db/history
news/db/history.dir
news/db/history.hash
news/db/history.index
news/db/newsgroups
news/etc/actsync.cfg
news/etc/actsync.ign
news/etc/buffindexed.conf
news/etc/control.ctl
news/etc/cycbuff.conf
news/etc/distrib.pats
news/etc/expire.ctl
news/etc/incoming.conf
news/etc/inn.conf
news/etc/innfeed.conf
news/etc/innreport.conf
news/etc/innwatch.ctl
news/etc/moderators
news/etc/motd.news
news/etc/news2mail.cf
news/etc/newsfeeds
news/etc/nnrpd.track
news/etc/nntpsend.ctl
news/etc/ovdb.conf
news/etc/overview.fmt
news/etc/passwd.nntp
news/etc/radius.conf
news/etc/readers.conf
news/etc/sasl.conf
news/etc/storage.conf
news/etc/subscriptions
news/include/clibrary.h
news/include/config.h
news/include/configdata.h
news/include/dbz.h
news/include/libinn.h
news/include/storage.h
news/lib/docheckgroups
news/lib/innreport_inn.pm
news/lib/innshellvars
news/lib/innshellvars.pl
news/lib/innshellvars.tcl
news/lib/libinn.a
news/lib/libstorage.a
share/doc/inn/HISTORY
share/doc/inn/INSTALL
share/doc/inn/LICENSE
@dirrm share/doc/inn
@dirrm news/spool/tmp
@dirrm news/spool/overview
@dirrm news/spool/outgoing
@dirrm news/spool/innfeed
@dirrm news/spool/incoming/bad
@dirrm news/spool/incoming
@dirrm news/spool/articles
@dirrm news/spool/archive
@dirrm news/spool
@dirrm news/run
@dirrm news/lib
@dirrm news/include
@dirrm news/etc
@dirrm news/db
@dirrm news/bin/rnews.libexec
@dirrm news/bin/filter
@dirrm news/bin/control
@dirrm news/bin/auth/resolv
@dirrm news/bin/auth/passwd
@dirrm news/bin/auth
@dirrm news/bin
@dirrm news
