# mail/postgrey/files/postgrey.in
#!/bin/sh
# $FreeBSD$
#

# PROVIDE: postgrey
# REQUIRE: LOGIN
# BEFORE: mail
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable postgrey:
#
# postgrey_enable="YES"
#
# See perldoc postgrey or run `postgrey --help` to see options. You can change
# the --dbdir path by setting postgrey_dbdir.  Specify additional options with
# the postgrey_options variable.
#

. /etc/rc.subr

name=postgrey
rcvar=`set_rcvar`

command=%%PREFIX%%/sbin/postgrey
extra_commands=reload

stop_postcmd=stop_postcmd

stop_postcmd()
{
  rm -f $pidfile
}

# set defaults

load_rc_config $name

postgrey_dbdir=${postgrey_dbdir:-"%%DBDIR%%"}
postgrey_enable=${postgrey_enable:-"NO"}
postgrey_options=${postgrey_options:-"--inet:10023"}
postgrey_pidfile=${postgrey_pidfile:-"/var/run/postgrey.pid"}
postgrey_flags=${postgrey_flags:-"-d --pidfile=${postgrey_pidfile} \
	--dbdir=${postgrey_dbdir} ${postgrey_options}"}

pidfile="${postgrey_pidfile}"
required_dirs=${postgrey_dbdir}

run_rc_command "$1"
