# ports-mgmt/packagekit/pkg-plist
bin/packagekit-bugreport.sh
bin/pk-debuginfo-install
bin/pkcon
bin/pkgenpack
bin/pkmon
@unexec if cmp -s %D/etc/PackageKit/CommandNotFound.conf %D/etc/PackageKit/CommandNotFound.conf.default; then rm -f %D/etc/PackageKit/CommandNotFound.conf; fi
etc/PackageKit/CommandNotFound.conf
@exec [ -f %B/CommandNotFound.conf ] || cp %B/%f %B/CommandNotFound.conf
@unexec if cmp -s %D/etc/PackageKit/PackageKit.conf %D/etc/PackageKit/PackageKit.conf.default; then rm -f %D/etc/PackageKit/PackageKit.conf; fi
etc/PackageKit/PackageKit.conf
@exec [ -f %B/PackageKit.conf ] || cp %B/%f %B/PackageKit.conf
etc/PackageKit/Vendor.conf
etc/PackageKit/events/post-transaction.d/README
etc/PackageKit/events/pre-transaction.d/README
etc/cron.daily/packagekit-background.cron
etc/dbus-1/system.d/org.freedesktop.PackageKit.conf
etc/profile.d/PackageKit.sh
etc/sysconfig/packagekit-background
include/PackageKit/packagekit-glib2/packagekit.h
include/PackageKit/packagekit-glib2/pk-bitfield.h
include/PackageKit/packagekit-glib2/pk-catalog.h
include/PackageKit/packagekit-glib2/pk-category.h
include/PackageKit/packagekit-glib2/pk-client-helper.h
include/PackageKit/packagekit-glib2/pk-client-sync.h
include/PackageKit/packagekit-glib2/pk-client.h
include/PackageKit/packagekit-glib2/pk-common.h
include/PackageKit/packagekit-glib2/pk-control-sync.h
include/PackageKit/packagekit-glib2/pk-control.h
include/PackageKit/packagekit-glib2/pk-desktop.h
include/PackageKit/packagekit-glib2/pk-details.h
include/PackageKit/packagekit-glib2/pk-distro-upgrade.h
include/PackageKit/packagekit-glib2/pk-enum-types.h
include/PackageKit/packagekit-glib2/pk-enum.h
include/PackageKit/packagekit-glib2/pk-error.h
include/PackageKit/packagekit-glib2/pk-eula-required.h
include/PackageKit/packagekit-glib2/pk-files.h
include/PackageKit/packagekit-glib2/pk-item-progress.h
include/PackageKit/packagekit-glib2/pk-media-change-required.h
include/PackageKit/packagekit-glib2/pk-message.h
include/PackageKit/packagekit-glib2/pk-package-id.h
include/PackageKit/packagekit-glib2/pk-package-ids.h
include/PackageKit/packagekit-glib2/pk-package-sack-sync.h
include/PackageKit/packagekit-glib2/pk-package-sack.h
include/PackageKit/packagekit-glib2/pk-package.h
include/PackageKit/packagekit-glib2/pk-progress.h
include/PackageKit/packagekit-glib2/pk-repo-detail.h
include/PackageKit/packagekit-glib2/pk-repo-signature-required.h
include/PackageKit/packagekit-glib2/pk-require-restart.h
include/PackageKit/packagekit-glib2/pk-results.h
include/PackageKit/packagekit-glib2/pk-service-pack.h
include/PackageKit/packagekit-glib2/pk-source.h
include/PackageKit/packagekit-glib2/pk-task-sync.h
include/PackageKit/packagekit-glib2/pk-task.h
include/PackageKit/packagekit-glib2/pk-transaction-list.h
include/PackageKit/packagekit-glib2/pk-transaction-past.h
include/PackageKit/packagekit-glib2/pk-update-detail.h
include/PackageKit/packagekit-glib2/pk-version.h
include/PackageKit/plugin/packagekit-plugin.h
include/PackageKit/plugin/pk-backend-job.h
include/PackageKit/plugin/pk-backend-spawn.h
include/PackageKit/plugin/pk-backend.h
include/PackageKit/plugin/pk-conf.h
include/PackageKit/plugin/pk-plugin.h
include/PackageKit/plugin/pk-shared.h
include/PackageKit/plugin/pk-transaction.h
lib/girepository-1.0/PackageKitGlib-1.0.typelib
lib/girepository-1.0/PackageKitPlugin-1.0.typelib
lib/libpackagekit-glib2.a
lib/libpackagekit-glib2.la
lib/libpackagekit-glib2.so
lib/libpackagekit-glib2.so.16
lib/packagekit-backend/libpk_backend_dummy.a
lib/packagekit-backend/libpk_backend_dummy.la
lib/packagekit-backend/libpk_backend_dummy.so
%%PORTS%%lib/packagekit-backend/libpk_backend_ports.a
%%PORTS%%lib/packagekit-backend/libpk_backend_ports.la
%%PORTS%%lib/packagekit-backend/libpk_backend_ports.so
lib/packagekit-backend/libpk_backend_test_fail.a
lib/packagekit-backend/libpk_backend_test_fail.la
lib/packagekit-backend/libpk_backend_test_fail.so
lib/packagekit-backend/libpk_backend_test_nop.a
lib/packagekit-backend/libpk_backend_test_nop.la
lib/packagekit-backend/libpk_backend_test_nop.so
lib/packagekit-backend/libpk_backend_test_spawn.a
lib/packagekit-backend/libpk_backend_test_spawn.la
lib/packagekit-backend/libpk_backend_test_spawn.so
lib/packagekit-backend/libpk_backend_test_succeed.a
lib/packagekit-backend/libpk_backend_test_succeed.la
lib/packagekit-backend/libpk_backend_test_succeed.so
lib/packagekit-backend/libpk_backend_test_thread.a
lib/packagekit-backend/libpk_backend_test_thread.la
lib/packagekit-backend/libpk_backend_test_thread.so
lib/packagekit-plugins/libpk_plugin-check-shared-libraries-in-use.a
lib/packagekit-plugins/libpk_plugin-check-shared-libraries-in-use.la
lib/packagekit-plugins/libpk_plugin-check-shared-libraries-in-use.so
lib/packagekit-plugins/libpk_plugin-clear-firmware-requests.a
lib/packagekit-plugins/libpk_plugin-clear-firmware-requests.la
lib/packagekit-plugins/libpk_plugin-clear-firmware-requests.so
lib/packagekit-plugins/libpk_plugin-no-update-process.a
lib/packagekit-plugins/libpk_plugin-no-update-process.la
lib/packagekit-plugins/libpk_plugin-no-update-process.so
lib/packagekit-plugins/libpk_plugin-require-restart.a
lib/packagekit-plugins/libpk_plugin-require-restart.la
lib/packagekit-plugins/libpk_plugin-require-restart.so
lib/packagekit-plugins/libpk_plugin-scan-desktop-files.a
lib/packagekit-plugins/libpk_plugin-scan-desktop-files.la
lib/packagekit-plugins/libpk_plugin-scan-desktop-files.so
lib/packagekit-plugins/libpk_plugin-systemd-updates.a
lib/packagekit-plugins/libpk_plugin-systemd-updates.la
lib/packagekit-plugins/libpk_plugin-systemd-updates.so
lib/packagekit-plugins/libpk_plugin-update-check-processes.a
lib/packagekit-plugins/libpk_plugin-update-check-processes.la
lib/packagekit-plugins/libpk_plugin-update-check-processes.so
lib/packagekit-plugins/libpk_plugin-update-package-cache.a
lib/packagekit-plugins/libpk_plugin-update-package-cache.la
lib/packagekit-plugins/libpk_plugin-update-package-cache.so
lib/packagekit-plugins/libpk_plugin_scripts.a
lib/packagekit-plugins/libpk_plugin_scripts.la
lib/packagekit-plugins/libpk_plugin_scripts.so
%%PYTHON_SITELIBDIR%%/packagekit/__init__.py
%%PYTHON_SITELIBDIR%%/packagekit/__init__.pyc
%%PYTHON_SITELIBDIR%%/packagekit/__init__.pyo
%%PYTHON_SITELIBDIR%%/packagekit/backend.py
%%PYTHON_SITELIBDIR%%/packagekit/backend.pyc
%%PYTHON_SITELIBDIR%%/packagekit/backend.pyo
%%PYTHON_SITELIBDIR%%/packagekit/enums.py
%%PYTHON_SITELIBDIR%%/packagekit/enums.pyc
%%PYTHON_SITELIBDIR%%/packagekit/enums.pyo
%%PYTHON_SITELIBDIR%%/packagekit/filter.py
%%PYTHON_SITELIBDIR%%/packagekit/filter.pyc
%%PYTHON_SITELIBDIR%%/packagekit/filter.pyo
%%PYTHON_SITELIBDIR%%/packagekit/misc.py
%%PYTHON_SITELIBDIR%%/packagekit/misc.pyc
%%PYTHON_SITELIBDIR%%/packagekit/misc.pyo
%%PYTHON_SITELIBDIR%%/packagekit/package.py
%%PYTHON_SITELIBDIR%%/packagekit/package.pyc
%%PYTHON_SITELIBDIR%%/packagekit/package.pyo
%%PYTHON_SITELIBDIR%%/packagekit/progress.py
%%PYTHON_SITELIBDIR%%/packagekit/progress.pyc
%%PYTHON_SITELIBDIR%%/packagekit/progress.pyo
libdata/pkgconfig/packagekit-glib2.pc
libdata/pkgconfig/packagekit-plugin.pc
libexec/packagekitd
libexec/pk-clear-offline-update
libexec/pk-command-not-found
libexec/pk-offline-update
libexec/pk-trigger-offline-update
%%PORTS%%share/PackageKit/helpers/ports/portsBackend.rb
%%PORTS%%share/PackageKit/helpers/ports/ruby_packagekit/enums.rb
%%PORTS%%share/PackageKit/helpers/ports/ruby_packagekit/prints.rb
share/PackageKit/helpers/test_spawn/search-name.sh
share/PackageKit/pk-upgrade-distro.sh
@exec mkdir -p /var/lib/PackageKit
@exec if [ ! -f /var/lib/PackageKit/%f ]; then cp -f %B/%f /var/lib/PackageKit; chmod o-r /var/lib/PackageKit/%f; fi
share/dbus-1/interfaces/org.freedesktop.PackageKit.Transaction.xml
share/dbus-1/interfaces/org.freedesktop.PackageKit.xml
share/dbus-1/system-services/org.freedesktop.PackageKit.service
share/gir-1.0/PackageKitGlib-1.0.gir
share/gir-1.0/PackageKitPlugin-1.0.gir
share/gtk-doc/html/PackageKit/PackageKit-pk-bitfield.html
share/gtk-doc/html/PackageKit/PackageKit-pk-client-sync.html
share/gtk-doc/html/PackageKit/PackageKit-pk-common.html
share/gtk-doc/html/PackageKit/PackageKit-pk-control-sync.html
share/gtk-doc/html/PackageKit/PackageKit-pk-enum.html
share/gtk-doc/html/PackageKit/PackageKit-pk-package-id.html
share/gtk-doc/html/PackageKit/PackageKit-pk-package-ids.html
share/gtk-doc/html/PackageKit/PackageKit-pk-task-sync.html
share/gtk-doc/html/PackageKit/PackageKit.devhelp2
share/gtk-doc/html/PackageKit/PackageKit.html
share/gtk-doc/html/PackageKit/PkCatalog.html
share/gtk-doc/html/PackageKit/PkCategory.html
share/gtk-doc/html/PackageKit/PkClient.html
share/gtk-doc/html/PackageKit/PkControl.html
share/gtk-doc/html/PackageKit/PkDesktop.html
share/gtk-doc/html/PackageKit/PkDetails.html
share/gtk-doc/html/PackageKit/PkDistroUpgrade.html
share/gtk-doc/html/PackageKit/PkError.html
share/gtk-doc/html/PackageKit/PkEulaRequired.html
share/gtk-doc/html/PackageKit/PkFiles.html
share/gtk-doc/html/PackageKit/PkMediaChangeRequired.html
share/gtk-doc/html/PackageKit/PkMessage.html
share/gtk-doc/html/PackageKit/PkPackage.html
share/gtk-doc/html/PackageKit/PkPackageSack.html
share/gtk-doc/html/PackageKit/PkProgress.html
share/gtk-doc/html/PackageKit/PkRepoDetail.html
share/gtk-doc/html/PackageKit/PkRepoSignatureRequired.html
share/gtk-doc/html/PackageKit/PkRequireRestart.html
share/gtk-doc/html/PackageKit/PkResults.html
share/gtk-doc/html/PackageKit/PkServicePack.html
share/gtk-doc/html/PackageKit/PkTask.html
share/gtk-doc/html/PackageKit/PkTransactionPast.html
share/gtk-doc/html/PackageKit/PkUpdateDetail.html
share/gtk-doc/html/PackageKit/Transaction.html
share/gtk-doc/html/PackageKit/api-reference.html
share/gtk-doc/html/PackageKit/backend-compiled.html
share/gtk-doc/html/PackageKit/backend-spawn.html
share/gtk-doc/html/PackageKit/backends-spawn-common.html
share/gtk-doc/html/PackageKit/backends-spawn-methods.html
share/gtk-doc/html/PackageKit/concepts.html
share/gtk-doc/html/PackageKit/config-main.html
share/gtk-doc/html/PackageKit/developer-faq.html
share/gtk-doc/html/PackageKit/faq.html
share/gtk-doc/html/PackageKit/home.png
share/gtk-doc/html/PackageKit/index.html
share/gtk-doc/html/PackageKit/index.sgml
share/gtk-doc/html/PackageKit/introduction-backends.html
share/gtk-doc/html/PackageKit/introduction-cancellation.html
share/gtk-doc/html/PackageKit/introduction-errors.html
share/gtk-doc/html/PackageKit/introduction-group-type.html
share/gtk-doc/html/PackageKit/introduction-ideas-filters.html
share/gtk-doc/html/PackageKit/introduction-ideas-status.html
share/gtk-doc/html/PackageKit/introduction-ideas-transactionid.html
share/gtk-doc/html/PackageKit/introduction-ideas-transactions.html
share/gtk-doc/html/PackageKit/introduction.html
share/gtk-doc/html/PackageKit/left.png
share/gtk-doc/html/PackageKit/libpackagekit-gobject.html
share/gtk-doc/html/PackageKit/libpackagekit-helpers.html
share/gtk-doc/html/PackageKit/pk-faq-error-code.png
share/gtk-doc/html/PackageKit/pk-faq-filter.png
share/gtk-doc/html/PackageKit/pk-faq-status.png
share/gtk-doc/html/PackageKit/pk-structure.png
share/gtk-doc/html/PackageKit/pk-structure.svg
share/gtk-doc/html/PackageKit/pk-transactions-auto-untrusted.png
share/gtk-doc/html/PackageKit/pk-transactions-download.png
share/gtk-doc/html/PackageKit/pk-transactions-failure.png
share/gtk-doc/html/PackageKit/pk-transactions-repair-required.png
share/gtk-doc/html/PackageKit/pk-transactions-set-locale.png
share/gtk-doc/html/PackageKit/pk-transactions-sig-install.png
share/gtk-doc/html/PackageKit/pk-transactions-success.png
share/gtk-doc/html/PackageKit/pk-transactions-trusted.png
share/gtk-doc/html/PackageKit/pk-transactions.svg
share/gtk-doc/html/PackageKit/right.png
share/gtk-doc/html/PackageKit/specification.html
share/gtk-doc/html/PackageKit/style.css
share/gtk-doc/html/PackageKit/up.png
share/locale/ar/LC_MESSAGES/PackageKit.mo
share/locale/as/LC_MESSAGES/PackageKit.mo
share/locale/bg/LC_MESSAGES/PackageKit.mo
share/locale/bg_BG/LC_MESSAGES/PackageKit.mo
share/locale/bn_IN/LC_MESSAGES/PackageKit.mo
share/locale/ca/LC_MESSAGES/PackageKit.mo
share/locale/ca@valencia/LC_MESSAGES/PackageKit.mo
share/locale/cs/LC_MESSAGES/PackageKit.mo
share/locale/da/LC_MESSAGES/PackageKit.mo
share/locale/de/LC_MESSAGES/PackageKit.mo
share/locale/el/LC_MESSAGES/PackageKit.mo
share/locale/en_GB/LC_MESSAGES/PackageKit.mo
share/locale/eo/LC_MESSAGES/PackageKit.mo
share/locale/es/LC_MESSAGES/PackageKit.mo
share/locale/eu/LC_MESSAGES/PackageKit.mo
share/locale/fa_IR/LC_MESSAGES/PackageKit.mo
share/locale/fi/LC_MESSAGES/PackageKit.mo
share/locale/fr/LC_MESSAGES/PackageKit.mo
share/locale/ga/LC_MESSAGES/PackageKit.mo
share/locale/gl/LC_MESSAGES/PackageKit.mo
share/locale/gu/LC_MESSAGES/PackageKit.mo
share/locale/he/LC_MESSAGES/PackageKit.mo
share/locale/hi/LC_MESSAGES/PackageKit.mo
share/locale/hu/LC_MESSAGES/PackageKit.mo
share/locale/id/LC_MESSAGES/PackageKit.mo
share/locale/it/LC_MESSAGES/PackageKit.mo
share/locale/it_IT/LC_MESSAGES/PackageKit.mo
share/locale/ja/LC_MESSAGES/PackageKit.mo
share/locale/kk/LC_MESSAGES/PackageKit.mo
share/locale/kn/LC_MESSAGES/PackageKit.mo
share/locale/ko/LC_MESSAGES/PackageKit.mo
share/locale/lt/LC_MESSAGES/PackageKit.mo
share/locale/lv/LC_MESSAGES/PackageKit.mo
share/locale/ml/LC_MESSAGES/PackageKit.mo
share/locale/mr/LC_MESSAGES/PackageKit.mo
share/locale/ms/LC_MESSAGES/PackageKit.mo
share/locale/nb/LC_MESSAGES/PackageKit.mo
share/locale/nl/LC_MESSAGES/PackageKit.mo
share/locale/or/LC_MESSAGES/PackageKit.mo
share/locale/pa/LC_MESSAGES/PackageKit.mo
share/locale/pl/LC_MESSAGES/PackageKit.mo
share/locale/pt/LC_MESSAGES/PackageKit.mo
share/locale/pt_BR/LC_MESSAGES/PackageKit.mo
share/locale/ro/LC_MESSAGES/PackageKit.mo
share/locale/ru/LC_MESSAGES/PackageKit.mo
share/locale/sk/LC_MESSAGES/PackageKit.mo
share/locale/sl/LC_MESSAGES/PackageKit.mo
share/locale/sr/LC_MESSAGES/PackageKit.mo
share/locale/sr@latin/LC_MESSAGES/PackageKit.mo
share/locale/sv/LC_MESSAGES/PackageKit.mo
share/locale/ta/LC_MESSAGES/PackageKit.mo
share/locale/te/LC_MESSAGES/PackageKit.mo
share/locale/th/LC_MESSAGES/PackageKit.mo
share/locale/tr/LC_MESSAGES/PackageKit.mo
share/locale/uk/LC_MESSAGES/PackageKit.mo
share/locale/wa/LC_MESSAGES/PackageKit.mo
share/locale/zh_CN/LC_MESSAGES/PackageKit.mo
share/locale/zh_TW/LC_MESSAGES/PackageKit.mo
share/mime/packages/packagekit-catalog.xml
share/mime/packages/packagekit-package-list.xml
share/mime/packages/packagekit-servicepack.xml
share/polkit-1/actions/org.freedesktop.packagekit.policy
@exec mkdir -p /var/cache/PackageKit/downloads
@unexec rmdir /var/cache/PackageKit/downloads 2>/dev/null || true
@unexec rmdir /var/cache/PackageKit 2>/dev/null || true
@unexec rmdir /var/cache 2>/dev/null || true
@unexec rmdir /var/lib/PackageKit 2>/dev/null || true
@unexec rmdir /var/lib 2>/dev/null || true
@dirrmtry share/locale/te/LC_MESSAGES
@dirrmtry share/locale/te
@dirrmtry share/locale/sr@latin/LC_MESSAGES
@dirrmtry share/locale/sr@latin
@dirrmtry share/locale/mr/LC_MESSAGES
@dirrmtry share/locale/mr
@dirrmtry share/locale/it_IT/LC_MESSAGES
@dirrmtry share/locale/it_IT
@dirrmtry share/locale/ca@valencia/LC_MESSAGES
@dirrmtry share/locale/ca@valencia
@dirrmtry share/locale/bn_IN/LC_MESSAGES
@dirrmtry share/locale/bn_IN
@dirrmtry share/locale/bg_BG/LC_MESSAGES
@dirrmtry share/locale/bg_BG
@dirrmtry share/locale/as/LC_MESSAGES
@dirrmtry share/locale/as
@dirrm share/gtk-doc/html/PackageKit
@dirrmtry share/gtk-doc/html
@dirrmtry share/gtk-doc
@dirrm share/PackageKit/helpers/test_spawn
%%PORTS%%@dirrm share/PackageKit/helpers/ports/ruby_packagekit
%%PORTS%%@dirrm share/PackageKit/helpers/ports
@dirrm share/PackageKit/helpers
@dirrm share/PackageKit
@dirrm %%PYTHON_SITELIBDIR%%/packagekit
@dirrm lib/packagekit-plugins
@dirrm lib/packagekit-backend
@dirrm include/PackageKit/packagekit-glib2
@dirrm include/PackageKit/plugin
@dirrm include/PackageKit
@dirrmtry etc/sysconfig
@dirrmtry etc/profile.d
@dirrmtry etc/cron.daily
@dirrmtry etc/PackageKit/events/post-transaction.d
@dirrmtry etc/PackageKit/events/pre-transaction.d
@dirrmtry etc/PackageKit/events
@dirrmtry etc/PackageKit
