# ports-mgmt/packagekit/pkg-plist
bin/pkcon
bin/pkmon
%%ETCDIR%%/PackageKit.conf
%%ETCDIR%%/Vendor.conf
etc/cron.daily/packagekit-background.cron
etc/sysconfig/packagekit-background
include/PackageKit/packagekit-glib2/packagekit.h
include/PackageKit/packagekit-glib2/pk-bitfield.h
include/PackageKit/packagekit-glib2/pk-category.h
include/PackageKit/packagekit-glib2/pk-client-helper.h
include/PackageKit/packagekit-glib2/pk-client-sync.h
include/PackageKit/packagekit-glib2/pk-client.h
include/PackageKit/packagekit-glib2/pk-common.h
include/PackageKit/packagekit-glib2/pk-control-sync.h
include/PackageKit/packagekit-glib2/pk-control.h
include/PackageKit/packagekit-glib2/pk-desktop.h
include/PackageKit/packagekit-glib2/pk-details.h
include/PackageKit/packagekit-glib2/pk-distro-upgrade.h
include/PackageKit/packagekit-glib2/pk-enum-types.h
include/PackageKit/packagekit-glib2/pk-enum.h
include/PackageKit/packagekit-glib2/pk-error.h
include/PackageKit/packagekit-glib2/pk-eula-required.h
include/PackageKit/packagekit-glib2/pk-files.h
include/PackageKit/packagekit-glib2/pk-item-progress.h
include/PackageKit/packagekit-glib2/pk-media-change-required.h
include/PackageKit/packagekit-glib2/pk-offline.h
include/PackageKit/packagekit-glib2/pk-package-id.h
include/PackageKit/packagekit-glib2/pk-package-ids.h
include/PackageKit/packagekit-glib2/pk-package-sack-sync.h
include/PackageKit/packagekit-glib2/pk-package-sack.h
include/PackageKit/packagekit-glib2/pk-package.h
include/PackageKit/packagekit-glib2/pk-progress.h
include/PackageKit/packagekit-glib2/pk-repo-detail.h
include/PackageKit/packagekit-glib2/pk-repo-signature-required.h
include/PackageKit/packagekit-glib2/pk-require-restart.h
include/PackageKit/packagekit-glib2/pk-results.h
include/PackageKit/packagekit-glib2/pk-source.h
include/PackageKit/packagekit-glib2/pk-task-sync.h
include/PackageKit/packagekit-glib2/pk-task.h
include/PackageKit/packagekit-glib2/pk-transaction-list.h
include/PackageKit/packagekit-glib2/pk-transaction-past.h
include/PackageKit/packagekit-glib2/pk-update-detail.h
include/PackageKit/packagekit-glib2/pk-version.h
lib/girepository-1.0/PackageKitGlib-1.0.typelib
lib/gnome-settings-daemon-3.0/gtk-modules/pk-gtk-module.desktop
lib/gtk-3.0/modules/libpk-gtk-module.so
lib/libpackagekit-glib2.so
lib/libpackagekit-glib2.so.18
lib/libpackagekit-glib2.so.18.1.3
lib/packagekit-backend/libpk_backend_dummy.so
lib/packagekit-backend/libpk_backend_freebsd.so
lib/packagekit-backend/libpk_backend_test_fail.so
lib/packagekit-backend/libpk_backend_test_nop.so
lib/packagekit-backend/libpk_backend_test_spawn.so
lib/packagekit-backend/libpk_backend_test_succeed.so
lib/packagekit-backend/libpk_backend_test_thread.so
%%PYTHON_SITELIBDIR%%/packagekit/__init__.py
%%PYTHON_SITELIBDIR%%/packagekit/backend.py
%%PYTHON_SITELIBDIR%%/packagekit/enums.py
%%PYTHON_SITELIBDIR%%/packagekit/filter.py
%%PYTHON_SITELIBDIR%%/packagekit/misc.py
%%PYTHON_SITELIBDIR%%/packagekit/package.py
%%PYTHON_SITELIBDIR%%/packagekit/progress.py
libdata/pkgconfig/packagekit-glib2.pc
libexec/packagekit-direct
libexec/packagekitd
%%DATADIR%%/helpers/test_spawn/search-name.sh
%%DATADIR%%/pk-upgrade-distro.sh
share/dbus-1/interfaces/org.freedesktop.PackageKit.Transaction.xml
share/dbus-1/interfaces/org.freedesktop.PackageKit.xml
share/dbus-1/system.d/org.freedesktop.PackageKit.conf
share/dbus-1/system-services/org.freedesktop.PackageKit.service
share/gir-1.0/PackageKitGlib-1.0.gir
share/locale/ar/LC_MESSAGES/PackageKit.mo
share/locale/as/LC_MESSAGES/PackageKit.mo
share/locale/ast/LC_MESSAGES/PackageKit.mo
share/locale/az/LC_MESSAGES/PackageKit.mo
share/locale/bg/LC_MESSAGES/PackageKit.mo
share/locale/bn_IN/LC_MESSAGES/PackageKit.mo
share/locale/br/LC_MESSAGES/PackageKit.mo
share/locale/ca/LC_MESSAGES/PackageKit.mo
share/locale/ca@valencia/LC_MESSAGES/PackageKit.mo
share/locale/cs/LC_MESSAGES/PackageKit.mo
share/locale/cy/LC_MESSAGES/PackageKit.mo
share/locale/da/LC_MESSAGES/PackageKit.mo
share/locale/de/LC_MESSAGES/PackageKit.mo
share/locale/el/LC_MESSAGES/PackageKit.mo
share/locale/en_GB/LC_MESSAGES/PackageKit.mo
share/locale/eo/LC_MESSAGES/PackageKit.mo
share/locale/es/LC_MESSAGES/PackageKit.mo
share/locale/eu/LC_MESSAGES/PackageKit.mo
share/locale/fa/LC_MESSAGES/PackageKit.mo
share/locale/fi/LC_MESSAGES/PackageKit.mo
share/locale/fo/LC_MESSAGES/PackageKit.mo
share/locale/fr/LC_MESSAGES/PackageKit.mo
share/locale/ga/LC_MESSAGES/PackageKit.mo
share/locale/gl/LC_MESSAGES/PackageKit.mo
share/locale/gu/LC_MESSAGES/PackageKit.mo
share/locale/he/LC_MESSAGES/PackageKit.mo
share/locale/hi/LC_MESSAGES/PackageKit.mo
share/locale/hr/LC_MESSAGES/PackageKit.mo
share/locale/hu/LC_MESSAGES/PackageKit.mo
share/locale/ia/LC_MESSAGES/PackageKit.mo
share/locale/id/LC_MESSAGES/PackageKit.mo
share/locale/it/LC_MESSAGES/PackageKit.mo
share/locale/ja/LC_MESSAGES/PackageKit.mo
share/locale/ka/LC_MESSAGES/PackageKit.mo
share/locale/kk/LC_MESSAGES/PackageKit.mo
share/locale/kn/LC_MESSAGES/PackageKit.mo
share/locale/ko/LC_MESSAGES/PackageKit.mo
share/locale/lt/LC_MESSAGES/PackageKit.mo
share/locale/lv/LC_MESSAGES/PackageKit.mo
share/locale/ml/LC_MESSAGES/PackageKit.mo
share/locale/mr/LC_MESSAGES/PackageKit.mo
share/locale/ms/LC_MESSAGES/PackageKit.mo
share/locale/nb/LC_MESSAGES/PackageKit.mo
share/locale/nl/LC_MESSAGES/PackageKit.mo
share/locale/nn/LC_MESSAGES/PackageKit.mo
share/locale/oc/LC_MESSAGES/PackageKit.mo
share/locale/or/LC_MESSAGES/PackageKit.mo
share/locale/pa/LC_MESSAGES/PackageKit.mo
share/locale/pl/LC_MESSAGES/PackageKit.mo
share/locale/pt/LC_MESSAGES/PackageKit.mo
share/locale/pt_BR/LC_MESSAGES/PackageKit.mo
share/locale/ro/LC_MESSAGES/PackageKit.mo
share/locale/ru/LC_MESSAGES/PackageKit.mo
share/locale/si/LC_MESSAGES/PackageKit.mo
share/locale/sk/LC_MESSAGES/PackageKit.mo
share/locale/sl/LC_MESSAGES/PackageKit.mo
share/locale/sq/LC_MESSAGES/PackageKit.mo
share/locale/sr/LC_MESSAGES/PackageKit.mo
share/locale/sr@latin/LC_MESSAGES/PackageKit.mo
share/locale/sv/LC_MESSAGES/PackageKit.mo
share/locale/ta/LC_MESSAGES/PackageKit.mo
share/locale/te/LC_MESSAGES/PackageKit.mo
share/locale/th/LC_MESSAGES/PackageKit.mo
share/locale/tr/LC_MESSAGES/PackageKit.mo
share/locale/uk/LC_MESSAGES/PackageKit.mo
share/locale/vi/LC_MESSAGES/PackageKit.mo
share/locale/wa/LC_MESSAGES/PackageKit.mo
share/locale/zh_CN/LC_MESSAGES/PackageKit.mo
share/locale/zh_HK/LC_MESSAGES/PackageKit.mo
share/locale/zh_TW/LC_MESSAGES/PackageKit.mo
share/man/man1/pkcon.1.gz
share/man/man1/pkmon.1.gz
share/metainfo/org.freedesktop.packagekit.metainfo.xml
share/polkit-1/actions/org.freedesktop.packagekit.policy
share/polkit-1/rules.d/org.freedesktop.packagekit.rules
share/vala/vapi/packagekit-glib2.deps
share/vala/vapi/packagekit-glib2.vapi
/var/lib/PackageKit/transactions.db
@dir /var/lib/PackageKit
@dir /var/lib
