# sysutils/mbmon/files/mbmon.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: mbmon
# REQUIRE: LOGIN
# BEFORE:  securelevel
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable the mbmon daemon:
#
# mbmon_enable="YES"
# mbmon_flags="<set as needed>"
#
# See mbmon(1) for mbmon_flags
# Set mbmon_port to the TCP port to listen to, default is 12999

. /etc/rc.subr

name="mbmon"
rcvar=mbmon_enable

# read configuration and set defaults
load_rc_config "$name"
: ${mbmon_enable="NO"}
: ${mbmon_port="12999"}

command="%%PREFIX%%/bin/${name}"
command_args="-P ${mbmon_port}"

run_rc_command "$1"
