# net/mpd/pkg-plist
etc/mpd/mpd.conf.sample
etc/mpd/mpd.links.sample
etc/mpd/mpd.secret.sample
etc/mpd/mpd.script.sample
etc/rc.d/mpd.sh
sbin/mpd
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/mpd.ps
%%PORTDOCS%%%%DOCSDIR%%/mpd.html
%%PORTDOCS%%%%DOCSDIR%%/mpd1.html
%%PORTDOCS%%%%DOCSDIR%%/mpd10.html
%%PORTDOCS%%%%DOCSDIR%%/mpd11.html
%%PORTDOCS%%%%DOCSDIR%%/mpd12.html
%%PORTDOCS%%%%DOCSDIR%%/mpd13.html
%%PORTDOCS%%%%DOCSDIR%%/mpd14.html
%%PORTDOCS%%%%DOCSDIR%%/mpd15.html
%%PORTDOCS%%%%DOCSDIR%%/mpd16.html
%%PORTDOCS%%%%DOCSDIR%%/mpd17.html
%%PORTDOCS%%%%DOCSDIR%%/mpd18.html
%%PORTDOCS%%%%DOCSDIR%%/mpd19.html
%%PORTDOCS%%%%DOCSDIR%%/mpd2.html
%%PORTDOCS%%%%DOCSDIR%%/mpd20.html
%%PORTDOCS%%%%DOCSDIR%%/mpd21.html
%%PORTDOCS%%%%DOCSDIR%%/mpd22.html
%%PORTDOCS%%%%DOCSDIR%%/mpd23.html
%%PORTDOCS%%%%DOCSDIR%%/mpd24.html
%%PORTDOCS%%%%DOCSDIR%%/mpd25.html
%%PORTDOCS%%%%DOCSDIR%%/mpd26.html
%%PORTDOCS%%%%DOCSDIR%%/mpd27.html
%%PORTDOCS%%%%DOCSDIR%%/mpd28.html
%%PORTDOCS%%%%DOCSDIR%%/mpd29.html
%%PORTDOCS%%%%DOCSDIR%%/mpd3.html
%%PORTDOCS%%%%DOCSDIR%%/mpd30.html
%%PORTDOCS%%%%DOCSDIR%%/mpd31.html
%%PORTDOCS%%%%DOCSDIR%%/mpd32.html
%%PORTDOCS%%%%DOCSDIR%%/mpd33.html
%%PORTDOCS%%%%DOCSDIR%%/mpd34.html
%%PORTDOCS%%%%DOCSDIR%%/mpd35.html
%%PORTDOCS%%%%DOCSDIR%%/mpd36.html
%%PORTDOCS%%%%DOCSDIR%%/mpd37.html
%%PORTDOCS%%%%DOCSDIR%%/mpd38.html
%%PORTDOCS%%%%DOCSDIR%%/mpd39.html
%%PORTDOCS%%%%DOCSDIR%%/mpd4.html
%%PORTDOCS%%%%DOCSDIR%%/mpd40.html
%%PORTDOCS%%%%DOCSDIR%%/mpd41.html
%%PORTDOCS%%%%DOCSDIR%%/mpd42.html
%%PORTDOCS%%%%DOCSDIR%%/mpd43.html
%%PORTDOCS%%%%DOCSDIR%%/mpd44.html
%%PORTDOCS%%%%DOCSDIR%%/mpd45.html
%%PORTDOCS%%%%DOCSDIR%%/mpd46.html
%%PORTDOCS%%%%DOCSDIR%%/mpd47.html
%%PORTDOCS%%%%DOCSDIR%%/mpd48.html
%%PORTDOCS%%%%DOCSDIR%%/mpd49.html
%%PORTDOCS%%%%DOCSDIR%%/mpd5.html
%%PORTDOCS%%%%DOCSDIR%%/mpd50.html
%%PORTDOCS%%%%DOCSDIR%%/mpd51.html
%%PORTDOCS%%%%DOCSDIR%%/mpd52.html
%%PORTDOCS%%%%DOCSDIR%%/mpd6.html
%%PORTDOCS%%%%DOCSDIR%%/mpd7.html
%%PORTDOCS%%%%DOCSDIR%%/mpd8.html
%%PORTDOCS%%%%DOCSDIR%%/mpd9.html
%%PORTDOCS%%%%DOCSDIR%%/mpd_toc.html
@dirrm %%PORTDOCS%%%%DOCSDIR%%
@dirrm etc/mpd
