# net/howl/files/openrc-mdnsresponder.in
#!/sbin/openrc-run
#
# Howl's mDNSResponder, a Zeroconf (Bonjour) service advertisement daemon.

name=mdnsresponder
command=%%PREFIX%%/bin/mDNSResponder
