# mail/dovecot-devel/pkg-plist
etc/dovecot-example.conf
etc/rc.d/dovecot-imapd.sh
lib/dovecot/imap
lib/dovecot/imap-auth
lib/dovecot/imap-login
sbin/imap-master
%%PORTDOCS%%share/doc/dovecot/auth.txt
%%PORTDOCS%%share/doc/dovecot/design.txt
%%PORTDOCS%%share/doc/dovecot/nfs.txt
%%PORTDOCS%%share/doc/dovecot/index.txt
%%PORTDOCS%%share/doc/dovecot/multiaccess.txt
%%PORTDOCS%%@dirrm share/doc/dovecot
@dirrm lib/dovecot
