# net/gini/pkg-plist
bin/ginicon
bin/ginifeed
etc/rc.d/gini.sh.sample
sbin/gini
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%@dirrm %%DOCSDIR%%
