# net/rwhois/pkg-plist
etc/rc.d/rwhoisd.sh
lib/rwhois/a.com/attribute_defs/asn.tmpl
lib/rwhois/a.com/attribute_defs/contact.tmpl
lib/rwhois/a.com/attribute_defs/domain.tmpl
lib/rwhois/a.com/attribute_defs/guardian.tmpl
lib/rwhois/a.com/attribute_defs/host.tmpl
lib/rwhois/a.com/attribute_defs/org.tmpl
lib/rwhois/a.com/attribute_defs/referral.tmpl
lib/rwhois/a.com/data/asn/asn.txt
lib/rwhois/a.com/data/contact/contact.txt
lib/rwhois/a.com/data/domain/domain.txt
lib/rwhois/a.com/data/guardian/guardian.txt
lib/rwhois/a.com/data/host/host.txt
lib/rwhois/a.com/data/org/org.txt
lib/rwhois/a.com/data/referral/referral.txt
lib/rwhois/a.com/schema
lib/rwhois/a.com/soa
lib/rwhois/bin/rwhois_deleter
lib/rwhois/bin/rwhois_indexer
lib/rwhois/bin/rwhois_repack
lib/rwhois/net-10.0.0.0-8/attribute_defs/contact.tmpl
lib/rwhois/net-10.0.0.0-8/attribute_defs/guardian.tmpl
lib/rwhois/net-10.0.0.0-8/attribute_defs/host.tmpl
lib/rwhois/net-10.0.0.0-8/attribute_defs/network.tmpl
lib/rwhois/net-10.0.0.0-8/attribute_defs/referral.tmpl
lib/rwhois/net-10.0.0.0-8/data/network/network.txt
lib/rwhois/net-10.0.0.0-8/data/referral/referral.txt
lib/rwhois/net-10.0.0.0-8/schema
lib/rwhois/net-10.0.0.0-8/soa
lib/rwhois/net-fd00:1234::-32/attribute_defs/host.tmpl
lib/rwhois/net-fd00:1234::-32/attribute_defs/referral.tmpl
lib/rwhois/net-fd00:1234::-32/attribute_defs/contact.tmpl
lib/rwhois/net-fd00:1234::-32/attribute_defs/guardian.tmpl
lib/rwhois/net-fd00:1234::-32/attribute_defs/network.tmpl
lib/rwhois/net-fd00:1234::-32/data/referral/referral.txt
lib/rwhois/net-fd00:1234::-32/data/network/network.txt
lib/rwhois/net-fd00:1234::-32/soa
lib/rwhois/net-fd00:1234::-32/schema
lib/rwhois/rwhoisd.allow
lib/rwhois/rwhoisd.auth_area
lib/rwhois/rwhoisd.conf
lib/rwhois/rwhoisd.deny
lib/rwhois/rwhoisd.dir
lib/rwhois/rwhoisd.root
lib/rwhois/rwhoisd.x.dir
lib/rwhois/sbin/rwhoisd
%%PORTDOCS%%%%DOCSDIR%%/INSTALL.html
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%%%DOCSDIR%%/UPGRADE
%%PORTDOCS%%%%DOCSDIR%%/operations_guide.html
%%PORTDOCS%%%%DOCSDIR%%/operations_guide.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc2167.txt
%%PORTDOCS%%%%DOCSDIR%%/security.html
%%PORTDOCS%%%%DOCSDIR%%/security.txt
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm lib/rwhois/sbin
@dirrm lib/rwhois/net-10.0.0.0-8/data/referral
@dirrm lib/rwhois/net-10.0.0.0-8/data/network
@dirrm lib/rwhois/net-10.0.0.0-8/data
@dirrm lib/rwhois/net-10.0.0.0-8/attribute_defs
@dirrm lib/rwhois/net-10.0.0.0-8
@dirrm lib/rwhois/net-fd00:1234::-32/attribute_defs
@dirrm lib/rwhois/net-fd00:1234::-32/data/referral
@dirrm lib/rwhois/net-fd00:1234::-32/data/network
@dirrm lib/rwhois/net-fd00:1234::-32/data
@dirrm lib/rwhois/net-fd00:1234::-32
@dirrm lib/rwhois/bin
@dirrm lib/rwhois/a.com/data/referral
@dirrm lib/rwhois/a.com/data/org
@dirrm lib/rwhois/a.com/data/host
@dirrm lib/rwhois/a.com/data/guardian
@dirrm lib/rwhois/a.com/data/domain
@dirrm lib/rwhois/a.com/data/contact
@dirrm lib/rwhois/a.com/data/asn
@dirrm lib/rwhois/a.com/data
@dirrm lib/rwhois/a.com/attribute_defs
@dirrm lib/rwhois/a.com
@dirrm lib/rwhois
