# net/opentracker/files/opentracker.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: opentracker
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf or /etc/rc.conf.local to
# enable opentracker:
# opentracker_enable (bool):	Set to "NO" by default.
# opentracker_config (path):	Set to 
#			"%%PREFIX%%/etc/opentracker/opentracker.conf"
#			by default

. /etc/rc.subr

name=opentracker
rcvar=opentracker_enable
desc="Lightweight Bittorrent tracker"

load_rc_config $name

: ${opentracker_enable:=NO}
: ${opentracker_config:="%%PREFIX%%/etc/opentracker/opentracker-ipv6.conf"}

command=/usr/sbin/daemon
pidfile=/var/run/opentracker.pid
command_args="-p $pidfile %%PREFIX%%/bin/opentracker -f ${opentracker_config}"
procname="%%PREFIX%%/bin/opentracker"
required_files=${opentracker_config}

run_rc_command $1
