# net/isc-dhcp44-relay/pkg-plist
%%IPV6%%etc/rc.d/isc-dhcrelay6
sbin/dhcrelay
share/man/man8/dhcrelay.8.gz
