# security/amavisd-new/files/amavis-milter.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: amavis_milter
# REQUIRE: amavisd
# BEFORE: mail
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable amavis-milter:
#
# amavis_milter_enable="YES"
#

. /etc/rc.subr

name=amavis_milter
rcvar=amavis_milter_enable

load_rc_config $name

: ${amavis_milter_enable:=NO}
: ${amavis_milter_flags=-D -p %%AMAVISDIR%%/amavis-milter.sock}
: ${amavis_milter_user=%%AMAVISUSER%%}

command=%%PREFIX%%/sbin/amavis-milter
required_files=%%PREFIX%%/etc/amavisd.conf

run_rc_command $1
