# japanese/rskkserv/files/rskkserv.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: rskkserv
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

prefix=%%PREFIX%%

rskkserv_enable=${rskkserv_enable:-"NO"}
rskkserv_flags=${rskkserv_flags:-""}

. /etc/rc.subr

name="rskkserv"
rcvar=rskkserv_enable
command="${prefix}/sbin/${name}"
command_interpreter="${prefix}/bin/ruby"
pidfile="/var/run/${name}.pid"
required_files="${prefix}/etc/${name}.conf"

load_rc_config $name
run_rc_command "$1"
