# finance/openerp/files/openerpd.in
#!/bin/sh

# PROVIDE: openerpd
# REQUIRE: postgresql LOGIN
#
# Add the following lines to /etc/rc.conf
# to enable this service:
#
# openerpd_enable (bool):   Set to NO by default.
#               Set it to YES to enable mongod.
#
# Additional configurable variables:
# openerpd_config (path):   Set to %%PREFIX%%/etc/openerp-server.conf
#               by default. Additional configuration. You
#               can also use openerpd_flags for additional
#               command line arguments.

. /etc/rc.subr

name="openerpd"
rcvar=${name}_enable

command=%%PREFIX%%/bin/openerp-server

load_rc_config $name

: ${openerpd_enable="NO"}
: ${openerpd_config="%%PREFIX%%/etc/openerp-server.conf"}

required_files="$openerpd_config"

command_args="-c $openerpd_config >/dev/null 2>&1 &"

run_rc_command "$1"
