# ftp/frox/pkg-plist
sbin/frox
etc/frox.conf.sample
etc/rc.d/frox.sh
%%PORTDOCS%%%%DOCSDIR%%/FAQ
%%PORTDOCS%%%%DOCSDIR%%/SECURITY
%%PORTDOCS%%%%DOCSDIR%%/CREDITS
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%@dirrm %%DOCSDIR%%
