# irc/dancer-services/pkg-plist
sbin/cleandb
sbin/dancer-services
sbin/encryptconf
sbin/encryptdb
sbin/fixlevel
sbin/servchk
sbin/shownicks
sbin/showchans
@unexec if cmp -s %D/etc/dancer-services/glines.conf %D/etc/dancer-services/glines.conf.sample; then rm -f %D/etc/dancer-services/glines.conf; fi
etc/dancer-services/glines.conf.sample
@exec if [ ! -f %B/glines.conf ]; then cp %B/%f %B/glines.conf; fi
@unexec if cmp -s %D/etc/dancer-services/jupes.conf %D/etc/dancer-services/jupes.conf.sample; then rm -f %D/etc/dancer-services/jupes.conf; fi
etc/dancer-services/jupes.conf.sample
@exec if [ ! -f %B/jupes.conf ]; then cp %B/%f %B/jupes.conf; fi
@unexec if cmp -s %D/etc/dancer-services/logon.news %D/etc/dancer-services/logon.news.sample; then rm -f %D/etc/dancer-services/logon.news; fi
etc/dancer-services/logon.news.sample
@exec if [ ! -f %B/logon.news ]; then cp %B/%f %B/logon.news; fi
@unexec if cmp -s %D/etc/dancer-services/motd.dcc %D/etc/dancer-services/motd.dcc.sample; then rm -f %D/etc/dancer-services/motd.dcc; fi
etc/dancer-services/motd.dcc.sample
@exec if [ ! -f %B/motd.dcc ]; then cp %B/%f %B/motd.dcc; fi
@unexec if cmp -s %D/etc/dancer-services/motd.global %D/etc/dancer-services/motd.global.sample; then rm -f %D/etc/dancer-services/motd.global; fi
etc/dancer-services/motd.global.sample
@exec if [ ! -f %B/motd.global ]; then cp %B/%f %B/motd.global; fi
@unexec if cmp -s %D/etc/dancer-services/services.conf %D/etc/dancer-services/services.conf.sample; then rm -f %D/etc/dancer-services/services.conf; fi
etc/dancer-services/services.conf.sample
@exec if [ ! -f %B/services.conf ]; then cp %B/%f %B/services.conf; fi
@unexec if cmp -s %D/etc/dancer-services/settings.conf %D/etc/dancer-services/settings.conf.sample; then rm -f %D/etc/dancer-services/settings.conf; fi
etc/dancer-services/settings.conf.sample
@exec if [ ! -f %B/settings.conf ]; then cp %B/%f %B/settings.conf; fi
etc/rc.d/dancer-services.sh
@exec install -d -o ircd -g ircd -m 700 /var/log/dancer-services /var/run/dancer-services
@unexec rmdir %D/etc/dancer-services /var/log/dancer-services /var/run/dancer-services 2>/dev/null || echo "If you are permanently removing this port, you should do a ``rm -rf ${PKG_PREFIX}/etc/dancer-services /var/log/dancer-services /var/run/dancer-services`` to remove any configuration files and log files left." | fmt
%%DATADIR%%/help/chanserv/access/add
%%DATADIR%%/help/chanserv/access/del
%%DATADIR%%/help/chanserv/access/index
%%DATADIR%%/help/chanserv/access/list
%%DATADIR%%/help/chanserv/akick/add
%%DATADIR%%/help/chanserv/akick/del
%%DATADIR%%/help/chanserv/akick/index
%%DATADIR%%/help/chanserv/akick/list
%%DATADIR%%/help/chanserv/autorem/add
%%DATADIR%%/help/chanserv/autorem/del
%%DATADIR%%/help/chanserv/autorem/index
%%DATADIR%%/help/chanserv/autorem/list
%%DATADIR%%/help/chanserv/clear
%%DATADIR%%/help/chanserv/clearnoexp
%%DATADIR%%/help/chanserv/drop
%%DATADIR%%/help/chanserv/forbid
%%DATADIR%%/help/chanserv/forget
%%DATADIR%%/help/chanserv/identify
%%DATADIR%%/help/chanserv/index
%%DATADIR%%/help/chanserv/info
%%DATADIR%%/help/chanserv/invite
%%DATADIR%%/help/chanserv/level/index
%%DATADIR%%/help/chanserv/level/list
%%DATADIR%%/help/chanserv/level/reset
%%DATADIR%%/help/chanserv/level/set
%%DATADIR%%/help/chanserv/list
%%DATADIR%%/help/chanserv/noexpire
%%DATADIR%%/help/chanserv/op
%%DATADIR%%/help/chanserv/register
%%DATADIR%%/help/chanserv/resetlevels
%%DATADIR%%/help/chanserv/set/email
%%DATADIR%%/help/chanserv/set/entrymsg
%%DATADIR%%/help/chanserv/set/forget
%%DATADIR%%/help/chanserv/set/founder
%%DATADIR%%/help/chanserv/set/guard
%%DATADIR%%/help/chanserv/set/index
%%DATADIR%%/help/chanserv/set/mlock
%%DATADIR%%/help/chanserv/set/password
%%DATADIR%%/help/chanserv/set/private
%%DATADIR%%/help/chanserv/set/restricted
%%DATADIR%%/help/chanserv/set/secure
%%DATADIR%%/help/chanserv/set/secureops
%%DATADIR%%/help/chanserv/set/splitops
%%DATADIR%%/help/chanserv/set/successor
%%DATADIR%%/help/chanserv/set/topiclock
%%DATADIR%%/help/chanserv/set/url
%%DATADIR%%/help/chanserv/set/verbose
%%DATADIR%%/help/chanserv/setpass
%%DATADIR%%/help/chanserv/status
%%DATADIR%%/help/chanserv/unban
%%DATADIR%%/help/chanserv/unforbid
%%DATADIR%%/help/chanserv/voice
%%DATADIR%%/help/global/gchannote
%%DATADIR%%/help/global/gcmemo
%%DATADIR%%/help/global/gmemo
%%DATADIR%%/help/global/gnote
%%DATADIR%%/help/global/index
%%DATADIR%%/help/global/motd
%%DATADIR%%/help/helpserv/chanserv
%%DATADIR%%/help/helpserv/index
%%DATADIR%%/help/helpserv/memoserv
%%DATADIR%%/help/helpserv/nickserv
%%DATADIR%%/help/helpserv/operserv
%%DATADIR%%/help/helpserv/statserv
%%DATADIR%%/help/memoserv/del
%%DATADIR%%/help/memoserv/forward
%%DATADIR%%/help/memoserv/index
%%DATADIR%%/help/memoserv/list
%%DATADIR%%/help/memoserv/purge
%%DATADIR%%/help/memoserv/read
%%DATADIR%%/help/memoserv/reply
%%DATADIR%%/help/memoserv/send
%%DATADIR%%/help/memoserv/undel
%%DATADIR%%/help/nickserv/access/add
%%DATADIR%%/help/nickserv/access/del
%%DATADIR%%/help/nickserv/access/index
%%DATADIR%%/help/nickserv/access/list
%%DATADIR%%/help/nickserv/collide
%%DATADIR%%/help/nickserv/drop
%%DATADIR%%/help/nickserv/droplink
%%DATADIR%%/help/nickserv/flag
%%DATADIR%%/help/nickserv/forbid
%%DATADIR%%/help/nickserv/ghost
%%DATADIR%%/help/nickserv/identify
%%DATADIR%%/help/nickserv/index
%%DATADIR%%/help/nickserv/info
%%DATADIR%%/help/nickserv/link
%%DATADIR%%/help/nickserv/list
%%DATADIR%%/help/nickserv/noexpire
%%DATADIR%%/help/nickserv/recover
%%DATADIR%%/help/nickserv/register
%%DATADIR%%/help/nickserv/release
%%DATADIR%%/help/nickserv/set/automask
%%DATADIR%%/help/nickserv/set/email
%%DATADIR%%/help/nickserv/set/gsm
%%DATADIR%%/help/nickserv/set/hide
%%DATADIR%%/help/nickserv/set/index
%%DATADIR%%/help/nickserv/set/kill
%%DATADIR%%/help/nickserv/set/master
%%DATADIR%%/help/nickserv/set/memos
%%DATADIR%%/help/nickserv/set/notify
%%DATADIR%%/help/nickserv/set/oper
%%DATADIR%%/help/nickserv/set/password
%%DATADIR%%/help/nickserv/set/phone
%%DATADIR%%/help/nickserv/set/private
%%DATADIR%%/help/nickserv/set/secure
%%DATADIR%%/help/nickserv/set/signon
%%DATADIR%%/help/nickserv/set/uin
%%DATADIR%%/help/nickserv/set/unsecure
%%DATADIR%%/help/nickserv/set/url
%%DATADIR%%/help/nickserv/setpass
%%DATADIR%%/help/nickserv/showlink
%%DATADIR%%/help/nickserv/unforbid
%%DATADIR%%/help/nickserv/unlink
%%DATADIR%%/help/operserv/channel
%%DATADIR%%/help/operserv/clones
%%DATADIR%%/help/operserv/dcc/boot
%%DATADIR%%/help/operserv/dcc/index
%%DATADIR%%/help/operserv/dcc/link
%%DATADIR%%/help/operserv/dcc/motd
%%DATADIR%%/help/operserv/dcc/quit
%%DATADIR%%/help/operserv/dcc/unlink
%%DATADIR%%/help/operserv/dcc/who
%%DATADIR%%/help/operserv/die
%%DATADIR%%/help/operserv/dump
%%DATADIR%%/help/operserv/floodoff
%%DATADIR%%/help/operserv/gline
%%DATADIR%%/help/operserv/help
%%DATADIR%%/help/operserv/htm
%%DATADIR%%/help/operserv/hub
%%DATADIR%%/help/operserv/identify
%%DATADIR%%/help/operserv/ignore
%%DATADIR%%/help/operserv/index
%%DATADIR%%/help/operserv/join
%%DATADIR%%/help/operserv/jump
%%DATADIR%%/help/operserv/jupe
%%DATADIR%%/help/operserv/kill
%%DATADIR%%/help/operserv/killchan
%%DATADIR%%/help/operserv/killhost
%%DATADIR%%/help/operserv/kline
%%DATADIR%%/help/operserv/off
%%DATADIR%%/help/operserv/omode
%%DATADIR%%/help/operserv/on
%%DATADIR%%/help/operserv/part
%%DATADIR%%/help/operserv/rehash
%%DATADIR%%/help/operserv/reload
%%DATADIR%%/help/operserv/restart
%%DATADIR%%/help/operserv/save
%%DATADIR%%/help/operserv/secure
%%DATADIR%%/help/operserv/set
%%DATADIR%%/help/operserv/stats
%%DATADIR%%/help/operserv/status
%%DATADIR%%/help/operserv/trace
%%DATADIR%%/help/operserv/umode
%%DATADIR%%/help/operserv/ungline
%%DATADIR%%/help/operserv/unjupe
%%DATADIR%%/help/seenserv/index
%%DATADIR%%/help/seenserv/seen
%%DATADIR%%/help/seenserv/seennick
%%DATADIR%%/help/seenserv/seenstat
%%DATADIR%%/help/seenserv/unseen
%%DATADIR%%/help/statserv/clearstats
%%DATADIR%%/help/statserv/domain
%%DATADIR%%/help/statserv/greplog
%%DATADIR%%/help/statserv/host
%%DATADIR%%/help/statserv/index
%%DATADIR%%/help/statserv/refresh
%%DATADIR%%/help/statserv/server
%%DATADIR%%/help/statserv/showadmins
%%DATADIR%%/help/statserv/showopers
%%DATADIR%%/help/statserv/showstats
%%DATADIR%%/help/statserv/split
%%DATADIR%%/help/statserv/stats
@dirrm %%DATADIR%%/help/statserv
@dirrm %%DATADIR%%/help/seenserv
@dirrm %%DATADIR%%/help/operserv/dcc
@dirrm %%DATADIR%%/help/operserv
@dirrm %%DATADIR%%/help/nickserv/set
@dirrm %%DATADIR%%/help/nickserv/access
@dirrm %%DATADIR%%/help/nickserv
@dirrm %%DATADIR%%/help/memoserv
@dirrm %%DATADIR%%/help/helpserv
@dirrm %%DATADIR%%/help/global
@dirrm %%DATADIR%%/help/chanserv/set
@dirrm %%DATADIR%%/help/chanserv/level
@dirrm %%DATADIR%%/help/chanserv/autorem
@dirrm %%DATADIR%%/help/chanserv/akick
@dirrm %%DATADIR%%/help/chanserv/access
@dirrm %%DATADIR%%/help/chanserv
@dirrm %%DATADIR%%/help
@dirrm %%DATADIR%%
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%@dirrm %%DOCSDIR%%
