# sysutils/docker/files/docker.in
#!/bin/sh

# PROVIDE: docker
# REQUIRE: DAEMON
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="docker"
rcvar="docker_enable"

stop_cmd="docker_stop"
start_cmd="docker_start"
command="%%PREFIX%%/bin/docker"

load_rc_config $name

: ${docker_enable=NO}
: ${docker_dir=/usr/docker}
: ${docker_nat_pf=YES}
: ${docker_nat_iface=NONE}

docker_start()
{
	if [ ! -d "$docker_dir" ] ; then
		echo "Missing $docker_dir! Please create / mount a ZFS dataset at this location."
		exit 1
	fi

	echo "Starting docker..."
	daemon -p /var/run/docker.pid ${command} -d -e jail -s zfs -g ${docker_dir} -D >/var/log/docker.log 2>/var/log/docker.log

	# Check for NAT support via PF
	# This is an ugly experimental hack for now, eventually will go away
	if [ "${docker_nat_pf}" != "YES" ] ; then return ; fi

	# Check if PF rules already loaded
	pfctl -s nat 2>/dev/null | grep -q 172.17
	if [ $? -eq 0 ] ; then return ; fi

	if [ "$docker_nat_iface" != "NONE" ] ; then
		iface="$docker_nat_iface"
	else
		iface=`netstat -f inet -nrW | grep '^default' | awk '{ print $6 }'`
	fi
	echo "nat on $iface from 172.17.0.0/16 to any -> ($iface)" > /tmp/pf-nat-docker.$$
	pfctl -f /tmp/pf-nat-docker.$$ 2>/dev/null
	pfctl -e 2>/dev/null
	rm /tmp/pf-nat-docker.$$
}

docker_stop()
{
	echo "Stopping docker..."
	if [ -e "/var/run/docker.pid" ] ; then
		pkill -F /var/run/docker.pid 
	fi
}

run_rc_command "$1"

