# games/kaid/pkg-plist
etc/kaid.conf.sample
etc/rc.d/kaid.sh
sbin/kaid
