# dns/dhisd/pkg-plist
sbin/dhisd
sbin/genkeys
@(,,600) etc/dhis/dhis.db.sample
@(,,600) etc/dhis/services.db.sample
etc/rc.d/dhisd.sh
@dir(,,700) etc/dhis
