# ports-mgmt/tinderbox-devel/pkg-plist
etc/rc.d/tinderbox
tinderbox/scripts/README
tinderbox/scripts/contrib/hooks/compress_wrkdir/README.txt
tinderbox/scripts/contrib/hooks/compress_wrkdir/compress_wrkdir.sh
tinderbox/scripts/contrib/hooks/kse-fix/README.txt
tinderbox/scripts/contrib/hooks/kse-fix/kse-fix.sh
tinderbox/scripts/contrib/hooks/kse-fix/libmap.conf
tinderbox/scripts/contrib/hooks/make-jobs-safe/README.txt
tinderbox/scripts/contrib/hooks/make-jobs-safe/make-jobs-safe.sh
tinderbox/scripts/contrib/hooks/mass-build-pgsql/mass-build-pgsql.sh
tinderbox/scripts/contrib/hooks/mass-build-pgsql/README.txt
tinderbox/scripts/ds.ph.dist
tinderbox/scripts/etc/env/.keep_me
tinderbox/scripts/etc/rc.d/tinderd
tinderbox/scripts/lib/Tinderbox/Build.pm
tinderbox/scripts/lib/Tinderbox/BuildPortsQueue.pm
tinderbox/scripts/lib/Tinderbox/Config.pm
tinderbox/scripts/lib/Tinderbox/Hook.pm
tinderbox/scripts/lib/Tinderbox/Jail.pm
tinderbox/scripts/lib/Tinderbox/MakeCache.pm
tinderbox/scripts/lib/Tinderbox/Port.pm
tinderbox/scripts/lib/Tinderbox/PortFailPattern.pm
tinderbox/scripts/lib/Tinderbox/PortFailReason.pm
tinderbox/scripts/lib/Tinderbox/PortsTree.pm
tinderbox/scripts/lib/Tinderbox/TinderObject.pm
tinderbox/scripts/lib/Tinderbox/TinderboxDS.pm
tinderbox/scripts/lib/Tinderbox/User.pm
tinderbox/scripts/lib/buildscript
tinderbox/scripts/lib/db-mysql.sh
tinderbox/scripts/lib/db-pgsql.sh
tinderbox/scripts/lib/enterbuild
tinderbox/scripts/lib/makemake
tinderbox/scripts/lib/pnohang.c
tinderbox/scripts/lib/portbuild
tinderbox/scripts/lib/tc_command.pl
tinderbox/scripts/lib/tc_command.sh
tinderbox/scripts/lib/tinderbox.env
tinderbox/scripts/lib/tinderlib.pl
tinderbox/scripts/lib/tinderlib.sh
tinderbox/scripts/sql/genschema
tinderbox/scripts/sql/schema.mysql.post
tinderbox/scripts/sql/schema.mysql.pre
tinderbox/scripts/sql/schema.pgsql.post
tinderbox/scripts/sql/schema.pgsql.pre
tinderbox/scripts/sql/values.config
tinderbox/scripts/sql/values.hooks
tinderbox/scripts/sql/values.lp
tinderbox/scripts/sql/values.pfp
tinderbox/scripts/sql/values.pfr
tinderbox/scripts/tc
tinderbox/scripts/tinderbox.ph.dist
tinderbox/scripts/tinderd
tinderbox/scripts/upgrade/build_ports.map
tinderbox/scripts/upgrade/build_ports_queue.map
tinderbox/scripts/upgrade/builds.map
tinderbox/scripts/upgrade/config.map
tinderbox/scripts/upgrade/hooks.map
tinderbox/scripts/upgrade/jails.map
tinderbox/scripts/upgrade/mig_pgsql_tinderbox-3.2.4_to_3.2.5.sql
tinderbox/scripts/upgrade/mig_mysql_tinderbox-3.0_to_3.1.sql
tinderbox/scripts/upgrade/mig_pgsql_tinderbox-3.0_to_3.1.sql
tinderbox/scripts/upgrade/mig_mysql_tinderbox-3.1_to_3.2.sql
tinderbox/scripts/upgrade/mig_pgsql_tinderbox-3.1_to_3.2.sql
tinderbox/scripts/upgrade/mig_mysql_tinderbox-3.2_to_3.2.1.sql
tinderbox/scripts/upgrade/mig_pgsql_tinderbox-3.2_to_3.2.1.sql
tinderbox/scripts/upgrade/mig_mysql_tinderbox-3.2.1_to_3.2.2.sql
tinderbox/scripts/upgrade/mig_pgsql_tinderbox-3.2.1_to_3.2.2.sql
tinderbox/scripts/upgrade/mig_mysql_tinderbox-3.2.2_to_3.2.3.sql
tinderbox/scripts/upgrade/mig_pgsql_tinderbox-3.2.2_to_3.2.3.sql
tinderbox/scripts/upgrade/mig_mysql_tinderbox-3.2.3_to_3.2.4.sql
tinderbox/scripts/upgrade/mig_pgsql_tinderbox-3.2.3_to_3.2.4.sql
tinderbox/scripts/upgrade/mig_mysql_tinderbox-3.2.4_to_3.2.5.sql
tinderbox/scripts/upgrade/mig_pgsql_tinderbox-3.2.5_to_3.2.6.sql
tinderbox/scripts/upgrade/mig_mysql_tinderbox-3.2.5_to_3.2.6.sql
tinderbox/scripts/upgrade/user_permissions.map
tinderbox/scripts/upgrade/order.lst
%%WEBUI%%tinderbox/scripts/webui/core/Build.php
%%WEBUI%%tinderbox/scripts/webui/core/BuildGroups.php
%%WEBUI%%tinderbox/scripts/webui/core/BuildPortsQueue.php
%%WEBUI%%tinderbox/scripts/webui/core/Config.php
%%WEBUI%%tinderbox/scripts/webui/core/Hooks.php
%%WEBUI%%tinderbox/scripts/webui/core/Jail.php
%%WEBUI%%tinderbox/scripts/webui/core/LogfilePattern.php
%%WEBUI%%tinderbox/scripts/webui/core/Port.php
%%WEBUI%%tinderbox/scripts/webui/core/PortFailPattern.php
%%WEBUI%%tinderbox/scripts/webui/core/PortFailReason.php
%%WEBUI%%tinderbox/scripts/webui/core/PortsTree.php
%%WEBUI%%tinderbox/scripts/webui/core/TinderObject.php
%%WEBUI%%tinderbox/scripts/webui/core/TinderboxDS.php
%%WEBUI%%tinderbox/scripts/webui/core/User.php
%%WEBUI%%tinderbox/scripts/webui/core/functions.php
%%WEBUI%%tinderbox/scripts/webui/favicon.ico
%%WEBUI%%tinderbox/scripts/webui/inc_ds.php.dist
%%WEBUI%%tinderbox/scripts/webui/inc_tinderbox.php.dist
%%WEBUI%%tinderbox/scripts/webui/index.php
%%WEBUI%%tinderbox/scripts/webui/module/module.php
%%WEBUI%%tinderbox/scripts/webui/module/moduleBuildGroups.php
%%WEBUI%%tinderbox/scripts/webui/module/moduleBuildPorts.php
%%WEBUI%%tinderbox/scripts/webui/module/moduleBuilds.php
%%WEBUI%%tinderbox/scripts/webui/module/moduleConfig.php
%%WEBUI%%tinderbox/scripts/webui/module/moduleLogs.php
%%WEBUI%%tinderbox/scripts/webui/module/modulePortFailureReasons.php
%%WEBUI%%tinderbox/scripts/webui/module/modulePorts.php
%%WEBUI%%tinderbox/scripts/webui/module/moduleRss.php
%%WEBUI%%tinderbox/scripts/webui/module/moduleSession.php
%%WEBUI%%tinderbox/scripts/webui/module/moduleTinderd.php
%%WEBUI%%tinderbox/scripts/webui/module/moduleUsers.php
%%WEBUI%%tinderbox/scripts/webui/templates/default/config.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/current_buildports.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/describe_port.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/display_build_groups.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/display_login.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/display_markup_log.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/failed_buildports.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/latest_buildports.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/list_buildports.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/list_builds.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/list_failure_reasons.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/list_tinderd_queue.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/messages.inc
%%WEBUI%%tinderbox/scripts/webui/templates/default/please_login.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/rss.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/tinderstyle.css
%%WEBUI%%tinderbox/scripts/webui/templates/default/tinderstyle.js
%%WEBUI%%tinderbox/scripts/webui/templates/default/user_admin.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/user_permissions.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/default/user_properties.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/config.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/current_buildports.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/describe_port.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/display_build_groups.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/display_login.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/display_markup_log.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/failed_buildports.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/footer.inc.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/header.inc.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/images/hdr_fill.png
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/latest_buildports.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/list_buildports.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/list_builds.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/list_failure_reasons.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/list_tinderd_queue.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/messages.inc
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/please_login.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/rss.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/tinderstyle.css
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/tinderstyle.js
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/user_admin.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/user_permissions.tpl
%%WEBUI%%tinderbox/scripts/webui/templates/paefchen/user_properties.tpl
%%WEBUI%%@dirrm tinderbox/scripts/webui/templates/paefchen/images
%%WEBUI%%@dirrm tinderbox/scripts/webui/templates/paefchen
%%WEBUI%%@dirrm tinderbox/scripts/webui/templates/default
%%WEBUI%%@dirrm tinderbox/scripts/webui/templates
%%WEBUI%%@dirrm tinderbox/scripts/webui/module
%%WEBUI%%@dirrm tinderbox/scripts/webui/core
%%WEBUI%%@dirrmtry tinderbox/scripts/webui
@dirrm tinderbox/scripts/upgrade
@dirrm tinderbox/scripts/sql
@dirrm tinderbox/scripts/lib/Tinderbox
@dirrm tinderbox/scripts/lib
@dirrm tinderbox/scripts/etc/rc.d
@dirrm tinderbox/scripts/etc/env
@dirrm tinderbox/scripts/etc
@dirrm tinderbox/scripts/contrib/hooks/mass-build-pgsql
@dirrm tinderbox/scripts/contrib/hooks/make-jobs-safe
@dirrm tinderbox/scripts/contrib/hooks/kse-fix
@dirrm tinderbox/scripts/contrib/hooks/compress_wrkdir
@dirrm tinderbox/scripts/contrib/hooks
@dirrm tinderbox/scripts/contrib
@dirrm tinderbox/scripts
@dirrmtry tinderbox
