# net-mgmt/chillispot/files/chillispot.sh.in
#!/bin/sh

# PROVIDE: chillispot
# REQUIRE: netif
# BEFORE:  pf ipfw ipfilter
# KEYWORD: nojail

. /etc/rc.subr

chillispot_enable=${chillispot_enable-"NO"}
chillispot_flags=${chillispot_flags-"--conf=/usr/local/etc/chilli.conf"}
chillispot_pidfile=${utility_pidfile-"/var/run/chilli.pid"}

. /etc/rc.subr

name="chillispot"
rcvar=`set_rcvar`
command="/usr/local/sbin/chilli"

load_rc_config $name

pidfile="${chillipot_pidfile}"

start_cmd="echo \"Starting ${name}.\"; /usr/bin/nice -5 ${command} ${chillispot_flags} ${command_args}"

run_rc_command "$1"
