# emulators/tpm-emulator/files/tpmd.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: tpmd
# REQUIRE: SERVERS
# BEFORE:  tcsd
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# tpmd_enable (bool):	Set to NO by default.
#			Set it to YES to enable tpmd.

. /etc/rc.subr

name=tpmd
rcvar="${name}_enable"
command="%%PREFIX%%/bin/${name}"

load_rc_config $name

: ${tpmd_enable:="NO"}
: ${tpmd_user:="%%USERS%%"}
: ${tpmd_group:="%%GROUPS%%"}

command_args="-o ${tpmd_user} -g ${tpmd_group}"
unset tpmd_user tpmd_group

run_rc_command "$1"
