# security/amavisd-milter/files/amavisd-milter.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: amavisd-milter
# REQUIRE: amavisd
# BEFORE: mail
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable amavisd-milter:
#
# amavisd_milter_enable="YES"
# amavisd_milter_flags="<set as needed>"
#

. /etc/rc.subr

name="amavisd_milter"
rcvar=amavisd_milter_enable

load_rc_config $name

# Set defaults
: ${amavisd_milter_enable:="NO"}
: ${amavisd_milter_user="%%AMAVISUSER%%"}
: ${amavisd_milter_socket:="%%AMAVISDIR%%/amavisd-milter.sock"}
: ${amavisd_milter_socket_perm:="0666"}

pidfile="${amavisd_milter_pidfile:-%%AMAVISDIR%%/amavisd-milter.pid}"
command="%%PREFIX%%/sbin/amavisd-milter"
command_args="-p ${pidfile} -s ${amavisd_milter_socket}"
required_files="%%PREFIX%%/etc/amavisd.conf"
start_postcmd="amavisd_milter_postcmd"

amavisd_milter_postcmd()
{
	chmod ${amavisd_milter_socket_perm} ${amavisd_milter_socket}
}

run_rc_command "$1"
