# net/pen/files/openrc-pen.in
#!/sbin/openrc-run

name=pen
command=%%PREFIX%%/bin/pen

depend() {
	use net
	keyword -shutdown
}
