# audio/daapd/pkg-plist
@comment $FreeBSD$
sbin/daapd
etc/daapd.conf.sample
etc/rc.d/daapd.sh
%%PORTDOCS%%%%DOCSDIR%%/README
@dirrm %%DATADIR%%
