# audio/daapd/pkg-plist
@comment $FreeBSD$
etc/daapd.conf.sample
etc/rc.d/daapd.sh
sbin/daapd
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%@dirrm %%DOCSDIR%%
