# net-mgmt/nfsen/files/nfsen.in
#! /bin/sh

# $FreeBSD$
#
# PROVIDE: nfsen
# REQUIRE: sshd
# KEYWORD: shutdown
#
nfsen_enable=${nfsen_enable-"NO"}
#nfsen_flags=${nfsen_flags-""}

. /etc/rc.subr

name="nfsen"
rcvar=`set_rcvar`
load_rc_config $name
command="%%PREFIX%%/bin/nfsen"
start_cmd="%%PREFIX%%/bin/nfsen start"
stop_cmd="%%PREFIX%%/bin/nfsen stop"
reload_cmd="%%PREFIX%%/bin/nfsen reconfig"
status_cmd="%%PREFIX%%/bin/nfsen status"

extra_commands="reload"

run_rc_command "$1"
