# ftp/bbftp/files/openrc-bbftpd.in
#!/sbin/openrc-run

name="bbftpd"
command=%%PREFIX%%/bin/${name}
pidfile=/var/run/${name}.pid
command_args="-b"
