# games/wolfpack/pkg-plist
wolfpack/README
wolfpack/bin/emp_client
wolfpack/bin/emp_server
wolfpack/bin/fairland
wolfpack/bin/files
wolfpack/bin/pconfig
wolfpack/data/auth
wolfpack/data/econfig
wolfpack/doc/CHANGES
wolfpack/doc/CHANGES2.3
wolfpack/doc/CHANGES3.0
wolfpack/doc/CHANGES4.X
wolfpack/doc/README
wolfpack/doc/WACKY_IDEAS
wolfpack/doc/WISH_LIST
wolfpack/doc/backups
wolfpack/doc/bug-reports
wolfpack/doc/coding
wolfpack/doc/compiling
wolfpack/doc/configuration
wolfpack/doc/debugging
wolfpack/doc/disk_space
wolfpack/doc/econfig
wolfpack/doc/files
wolfpack/doc/message
wolfpack/doc/pthreads
wolfpack/doc/threads
wolfpack/doc/troubleshooting
wolfpack/info.nr/Attack
wolfpack/info.nr/Autofeed
wolfpack/info.nr/Autonav
wolfpack/info.nr/BTU
wolfpack/info.nr/Basics
wolfpack/info.nr/Bridges
wolfpack/info.nr/Bugs
wolfpack/info.nr/CHANGES.CHAINSAW
wolfpack/info.nr/Citizens
wolfpack/info.nr/Clients
wolfpack/info.nr/Combat
wolfpack/info.nr/Commodities
wolfpack/info.nr/Communication
wolfpack/info.nr/Damage
wolfpack/info.nr/Deity
wolfpack/info.nr/Detection
wolfpack/info.nr/Diplomacy
wolfpack/info.nr/Distribution
wolfpack/info.nr/Education
wolfpack/info.nr/Empire2
wolfpack/info.nr/Empire3
wolfpack/info.nr/Empire4
wolfpack/info.nr/Expert
wolfpack/info.nr/Fallout
wolfpack/info.nr/Flak
wolfpack/info.nr/Food
wolfpack/info.nr/God
wolfpack/info.nr/Guerrilla
wolfpack/info.nr/Guide
wolfpack/info.nr/Happiness
wolfpack/info.nr/Hidden
wolfpack/info.nr/Hints
wolfpack/info.nr/History
wolfpack/info.nr/Hitchance
wolfpack/info.nr/Infrastructure
wolfpack/info.nr/Innards
wolfpack/info.nr/Interception
wolfpack/info.nr/Introduction
wolfpack/info.nr/Kill
wolfpack/info.nr/Land-units
wolfpack/info.nr/LandUnits
wolfpack/info.nr/Loans
wolfpack/info.nr/Maintenance
wolfpack/info.nr/Maps
wolfpack/info.nr/Market
wolfpack/info.nr/Merc
wolfpack/info.nr/Missions
wolfpack/info.nr/Mobility
wolfpack/info.nr/Move
wolfpack/info.nr/Moving
wolfpack/info.nr/Nation
wolfpack/info.nr/Novice
wolfpack/info.nr/Nuke-types
wolfpack/info.nr/Obsolete
wolfpack/info.nr/Occupation
wolfpack/info.nr/Old-empire
wolfpack/info.nr/Options
wolfpack/info.nr/Overview
wolfpack/info.nr/Plague
wolfpack/info.nr/Plane-types
wolfpack/info.nr/Planes
wolfpack/info.nr/Playing
wolfpack/info.nr/Populace
wolfpack/info.nr/Production
wolfpack/info.nr/Products
wolfpack/info.nr/Quick-ref
wolfpack/info.nr/Research
wolfpack/info.nr/Sector-ownership
wolfpack/info.nr/Sector-types
wolfpack/info.nr/Sectors
wolfpack/info.nr/Selector
wolfpack/info.nr/Server
wolfpack/info.nr/Ship-types
wolfpack/info.nr/Ships
wolfpack/info.nr/Spies
wolfpack/info.nr/Suggestions
wolfpack/info.nr/Syntax
wolfpack/info.nr/TOP
wolfpack/info.nr/Taxes
wolfpack/info.nr/Technology
wolfpack/info.nr/Time
wolfpack/info.nr/Trade-ships
wolfpack/info.nr/Unit-types
wolfpack/info.nr/Update
wolfpack/info.nr/Update-sequence
wolfpack/info.nr/Wolfpack
wolfpack/info.nr/accept
wolfpack/info.nr/add
wolfpack/info.nr/all
wolfpack/info.nr/announce
wolfpack/info.nr/anti
wolfpack/info.nr/apropos
wolfpack/info.nr/arm
wolfpack/info.nr/army
wolfpack/info.nr/assault
wolfpack/info.nr/attack
wolfpack/info.nr/bdes
wolfpack/info.nr/bestpath
wolfpack/info.nr/bmap
wolfpack/info.nr/board
wolfpack/info.nr/bomb
wolfpack/info.nr/break
wolfpack/info.nr/budget
wolfpack/info.nr/build
wolfpack/info.nr/buy
wolfpack/info.nr/bye
wolfpack/info.nr/capital
wolfpack/info.nr/cargo
wolfpack/info.nr/census
wolfpack/info.nr/change
wolfpack/info.nr/coastwatch
wolfpack/info.nr/collect
wolfpack/info.nr/commodity
wolfpack/info.nr/consider
wolfpack/info.nr/convert
wolfpack/info.nr/country
wolfpack/info.nr/cutoff
wolfpack/info.nr/declare
wolfpack/info.nr/deliver
wolfpack/info.nr/demobilize
wolfpack/info.nr/designate
wolfpack/info.nr/disable
wolfpack/info.nr/disarm
wolfpack/info.nr/distribute
wolfpack/info.nr/drop
wolfpack/info.nr/dump
wolfpack/info.nr/echo
wolfpack/info.nr/edit
wolfpack/info.nr/enable
wolfpack/info.nr/enlist
wolfpack/info.nr/execute
wolfpack/info.nr/explore
wolfpack/info.nr/financial
wolfpack/info.nr/fire
wolfpack/info.nr/flash
wolfpack/info.nr/fleetadd
wolfpack/info.nr/fly
wolfpack/info.nr/follow
wolfpack/info.nr/force
wolfpack/info.nr/fortify
wolfpack/info.nr/fuel
wolfpack/info.nr/give
wolfpack/info.nr/grind
wolfpack/info.nr/harden
wolfpack/info.nr/headlines
wolfpack/info.nr/help
wolfpack/info.nr/hidden
wolfpack/info.nr/improve
wolfpack/info.nr/info
wolfpack/info.nr/land
wolfpack/info.nr/launch
wolfpack/info.nr/lbmap
wolfpack/info.nr/lboard
wolfpack/info.nr/lcargo
wolfpack/info.nr/ldump
wolfpack/info.nr/ledger
wolfpack/info.nr/level
wolfpack/info.nr/list
wolfpack/info.nr/lload
wolfpack/info.nr/llookout
wolfpack/info.nr/lmap
wolfpack/info.nr/lmine
wolfpack/info.nr/load
wolfpack/info.nr/lookout
wolfpack/info.nr/lost
wolfpack/info.nr/lradar
wolfpack/info.nr/lrange
wolfpack/info.nr/lretreat
wolfpack/info.nr/lstat
wolfpack/info.nr/ltend
wolfpack/info.nr/lunload
wolfpack/info.nr/map
wolfpack/info.nr/march
wolfpack/info.nr/market
wolfpack/info.nr/mine
wolfpack/info.nr/mission
wolfpack/info.nr/mobupdate
wolfpack/info.nr/morale
wolfpack/info.nr/motd
wolfpack/info.nr/move
wolfpack/info.nr/mquota
wolfpack/info.nr/name
wolfpack/info.nr/nation
wolfpack/info.nr/navigate
wolfpack/info.nr/ndump
wolfpack/info.nr/newcap
wolfpack/info.nr/neweff
wolfpack/info.nr/newspaper
wolfpack/info.nr/nmap
wolfpack/info.nr/nuke
wolfpack/info.nr/offer
wolfpack/info.nr/offset
wolfpack/info.nr/order
wolfpack/info.nr/origin
wolfpack/info.nr/paradrop
wolfpack/info.nr/path
wolfpack/info.nr/payoff
wolfpack/info.nr/pbmap
wolfpack/info.nr/pboard
wolfpack/info.nr/pdump
wolfpack/info.nr/plane
wolfpack/info.nr/players
wolfpack/info.nr/pmap
wolfpack/info.nr/power
wolfpack/info.nr/pray
wolfpack/info.nr/production
wolfpack/info.nr/pstat
wolfpack/info.nr/qorder
wolfpack/info.nr/quit
wolfpack/info.nr/radar
wolfpack/info.nr/range
wolfpack/info.nr/read
wolfpack/info.nr/realm
wolfpack/info.nr/recon
wolfpack/info.nr/reject
wolfpack/info.nr/relations
wolfpack/info.nr/repay
wolfpack/info.nr/report
wolfpack/info.nr/reset
wolfpack/info.nr/resource
wolfpack/info.nr/retreat
wolfpack/info.nr/route
wolfpack/info.nr/sabotage
wolfpack/info.nr/sail
wolfpack/info.nr/satellite
wolfpack/info.nr/sbmap
wolfpack/info.nr/scrap
wolfpack/info.nr/scuttle
wolfpack/info.nr/sdump
wolfpack/info.nr/sect
wolfpack/info.nr/sell
wolfpack/info.nr/set
wolfpack/info.nr/setresource
wolfpack/info.nr/setsector
wolfpack/info.nr/sharebmap
wolfpack/info.nr/shark
wolfpack/info.nr/ship
wolfpack/info.nr/shoot
wolfpack/info.nr/show
wolfpack/info.nr/shutdown
wolfpack/info.nr/sinfrastructure
wolfpack/info.nr/skywatch
wolfpack/info.nr/smap
wolfpack/info.nr/sonar
wolfpack/info.nr/sorder
wolfpack/info.nr/spy
wolfpack/info.nr/sstat
wolfpack/info.nr/start
wolfpack/info.nr/starvation
wolfpack/info.nr/stop
wolfpack/info.nr/strength
wolfpack/info.nr/supply
wolfpack/info.nr/survey
wolfpack/info.nr/swapsector
wolfpack/info.nr/sweep
wolfpack/info.nr/telegram
wolfpack/info.nr/tend
wolfpack/info.nr/territory
wolfpack/info.nr/test
wolfpack/info.nr/threshold
wolfpack/info.nr/toggle
wolfpack/info.nr/torpedo
wolfpack/info.nr/trade
wolfpack/info.nr/transport
wolfpack/info.nr/treaty
wolfpack/info.nr/turn
wolfpack/info.nr/unload
wolfpack/info.nr/unsail
wolfpack/info.nr/update
wolfpack/info.nr/upgrade
wolfpack/info.nr/version
wolfpack/info.nr/wait
wolfpack/info.nr/wall
wolfpack/info.nr/wingadd
wolfpack/info.nr/wipe
wolfpack/info.nr/wire
wolfpack/info.nr/work
wolfpack/info.nr/zdone
wolfpack/scripts/README
wolfpack/scripts/bigstart
wolfpack/scripts/blitz/README
wolfpack/scripts/blitz/crontab.config
wolfpack/scripts/blitz/motd
wolfpack/scripts/blitz/setup-blitz
wolfpack/scripts/blitz/start-blitz
wolfpack/scripts/crontab.config
wolfpack/scripts/keepitup
wolfpack/scripts/mapper/Makefile
wolfpack/scripts/mapper/README
wolfpack/scripts/mapper/archive.pl
wolfpack/scripts/mapper/cleanup
wolfpack/scripts/mapper/crontab.config
wolfpack/scripts/mapper/default.param
wolfpack/scripts/mapper/makeframe
wolfpack/scripts/mapper/makempeg
wolfpack/scripts/mapper/map.param
wolfpack/scripts/mapper/mapper
wolfpack/scripts/mapper/over.awk
wolfpack/scripts/mapper/over2.pbm
wolfpack/scripts/mapper/param.pl
wolfpack/scripts/mapper/pgmmap.awk
wolfpack/scripts/mapper/randppm.awk
wolfpack/scripts/mapper/text.pbm
wolfpack/scripts/mapper/text.size
wolfpack/scripts/postresults
wolfpack/scripts/saveass
etc/rc.d/wolfpack.sh
@dirrm wolfpack/scripts/mapper
@dirrm wolfpack/scripts/blitz
@dirrm wolfpack/scripts
@dirrm wolfpack/info.nr
@dirrm wolfpack/doc
@dirrm wolfpack/data
@dirrm wolfpack/bin
@dirrm wolfpack
