# net-im/jabberd/pkg-plist
@unexec [ ! -f /var/jabberd/pid/router.pid ] || %D/etc/rc.d/jabberd stop || true
bin/c2s
bin/jabberd
bin/router
bin/s2s
bin/sm
@sample etc/jabberd/c2s.xml.sample
@sample etc/jabberd/jabberd.cfg.sample
@sample etc/jabberd/router-filter.xml.sample
@sample etc/jabberd/router-users.xml.sample
@sample etc/jabberd/router.xml.sample
@sample etc/jabberd/s2s.xml.sample
@sample etc/jabberd/sm.xml.sample
@sample etc/jabberd/templates/roster.xml.sample
%%ANON%%lib/jabberd/authreg_anon.so
%%BDB%%lib/jabberd/authreg_db.so
%%LDAP%%lib/jabberd/authreg_ldap.so
%%LDAP%%lib/jabberd/authreg_ldapfull.so
%%MYSQL%%lib/jabberd/authreg_mysql.so
%%PAM%%lib/jabberd/authreg_pam.so
%%PIPE%%lib/jabberd/authreg_pipe.so
%%PGSQL%%lib/jabberd/authreg_pgsql.so
%%SQLITE%%lib/jabberd/authreg_sqlite.so
lib/jabberd/libstorage.so
lib/jabberd/libstorage.so.0
lib/jabberd/libstorage.so.0.0.0
lib/jabberd/mod_active.so
lib/jabberd/mod_active.so.0
lib/jabberd/mod_active.so.0.0.0
lib/jabberd/mod_amp.so
lib/jabberd/mod_amp.so.0
lib/jabberd/mod_amp.so.0.0.0
lib/jabberd/mod_announce.so
lib/jabberd/mod_announce.so.0
lib/jabberd/mod_announce.so.0.0.0
lib/jabberd/mod_deliver.so
lib/jabberd/mod_deliver.so.0
lib/jabberd/mod_deliver.so.0.0.0
lib/jabberd/mod_disco.so
lib/jabberd/mod_disco.so.0
lib/jabberd/mod_disco.so.0.0.0
lib/jabberd/mod_echo.so
lib/jabberd/mod_echo.so.0
lib/jabberd/mod_echo.so.0.0.0
lib/jabberd/mod_help.so
lib/jabberd/mod_help.so.0
lib/jabberd/mod_help.so.0.0.0
lib/jabberd/mod_iq-last.so
lib/jabberd/mod_iq-last.so.0
lib/jabberd/mod_iq-last.so.0.0.0
lib/jabberd/mod_iq-ping.so
lib/jabberd/mod_iq-ping.so.0
lib/jabberd/mod_iq-ping.so.0.0.0
lib/jabberd/mod_iq-private.so
lib/jabberd/mod_iq-private.so.0
lib/jabberd/mod_iq-private.so.0.0.0
lib/jabberd/mod_iq-time.so
lib/jabberd/mod_iq-time.so.0
lib/jabberd/mod_iq-time.so.0.0.0
lib/jabberd/mod_iq-vcard.so
lib/jabberd/mod_iq-vcard.so.0
lib/jabberd/mod_iq-vcard.so.0.0.0
lib/jabberd/mod_iq-version.so
lib/jabberd/mod_iq-version.so.0
lib/jabberd/mod_iq-version.so.0.0.0
lib/jabberd/mod_offline.so
lib/jabberd/mod_offline.so.0
lib/jabberd/mod_offline.so.0.0.0
lib/jabberd/mod_pep.so
lib/jabberd/mod_pep.so.0
lib/jabberd/mod_pep.so.0.0.0
lib/jabberd/mod_presence.so
lib/jabberd/mod_presence.so.0
lib/jabberd/mod_presence.so.0.0.0
lib/jabberd/mod_privacy.so
lib/jabberd/mod_privacy.so.0
lib/jabberd/mod_privacy.so.0.0.0
lib/jabberd/mod_roster.so
lib/jabberd/mod_roster.so.0
lib/jabberd/mod_roster.so.0.0.0
lib/jabberd/mod_roster-publish.so
lib/jabberd/mod_roster-publish.so.0
lib/jabberd/mod_roster-publish.so.0.0.0
lib/jabberd/mod_session.so
lib/jabberd/mod_session.so.0
lib/jabberd/mod_session.so.0.0.0
lib/jabberd/mod_status.so
lib/jabberd/mod_status.so.0
lib/jabberd/mod_status.so.0.0.0
lib/jabberd/mod_template-roster.so
lib/jabberd/mod_template-roster.so.0
lib/jabberd/mod_template-roster.so.0.0.0
lib/jabberd/mod_vacation.so
lib/jabberd/mod_vacation.so.0
lib/jabberd/mod_vacation.so.0.0.0
lib/jabberd/mod_validate.so
lib/jabberd/mod_validate.so.0
lib/jabberd/mod_validate.so.0.0.0
lib/jabberd/mod_verify.so
lib/jabberd/mod_verify.so.0
lib/jabberd/mod_verify.so.0.0.0
%%FS%%lib/jabberd/storage_fs.so
%%BDB%%lib/jabberd/storage_db.so
%%LDAP%%lib/jabberd/storage_ldapvcard.so
%%MYSQL%%lib/jabberd/storage_mysql.so
%%PGSQL%%lib/jabberd/storage_pgsql.so
%%SQLITE%%lib/jabberd/storage_sqlite.so
man/man8/c2s.8.gz
man/man8/jabberd.8.gz
man/man8/router.8.gz
man/man8/s2s.8.gz
man/man8/sm.8.gz
@dir(%%JABBER_USER%%,%%JABBER_GROUP%%,0700) %%JABBER_RUNDIR%%/db
@dir(%%JABBER_USER%%,%%JABBER_GROUP%%,0700) %%JABBER_RUNDIR%%/logs
@dir(%%JABBER_USER%%,%%JABBER_GROUP%%,0700) %%JABBER_RUNDIR%%/pid
@dir(%%JABBER_USER%%,%%JABBER_GROUP%%,0700) %%JABBER_RUNDIR%%
