# mail/qpopper/qpopper.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: qpopper
# REQUIRE: SERVERS
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable qpopper:
# qpopper_enable="YES"

. /etc/rc.subr

name=qpopper
rcvar=qpopper_enable

load_rc_config $name

qpopper_enable=${qpopper_enable:-"NO"}

command=%%PREFIX%%/libexec/qpopper

run_rc_command "$1"
