# mail/qpopper/qpopper.in
#!/bin/sh
#
# PROVIDE: qpopper
# REQUIRE: SERVERS
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable qpopper:
# qpopper_enable="YES"
#

. /etc/rc.subr

qpopper_enable=${qpopper_enable-"NO"}
qpopper_flags=${qpopper_flags-""}

name=qpopper
rcvar=qpopper_enable
command=%%PREFIX%%/libexec/qpopper
command_args="${qpopper_flags}"

load_rc_config $name
run_rc_command "$1"
