# games/bnetd/pkg-plist
bin/bnclient
bin/bnpass
sbin/bnetd
etc/rc.d/bnetd.sh
@unexec rm -rf %D/bnetd 2>/dev/null || true
