# net/icradius/pkg-plist
etc/raddb/radius.conf.sample
etc/rc.d/icradiusd.sh
sbin/radiusd
sbin/checkrad
sbin/radwatch
%%CGI%%libexec/radius.cgi
%%CGI%%libexec/usage.cgi
%%PORTDOCS%%share/doc/icradius/COPYING
%%PORTDOCS%%share/doc/icradius/COPYRIGHT.Cistron
%%PORTDOCS%%share/doc/icradius/COPYRIGHT.ICRADIUS
%%PORTDOCS%%share/doc/icradius/COPYRIGHT.Livingston
%%PORTDOCS%%share/doc/icradius/README
%%PORTDOCS%%share/doc/icradius/README.Y2K
%%PORTDOCS%%share/doc/icradius/README.cisco
%%PORTDOCS%%share/doc/icradius/README.hints
%%PORTDOCS%%share/doc/icradius/README.proxy
%%PORTDOCS%%share/doc/icradius/README.simul
%%PORTDOCS%%share/doc/icradius/FAQ
%%PORTDOCS%%share/doc/icradius/THANKS
%%PORTDOCS%%share/doc/icradius/TODO
share/icradius/raddb/dictionary
share/icradius/raddb/dictionary.ascend
share/icradius/raddb/dictionary.cisco
share/icradius/raddb/dictionary.compat
share/icradius/raddb/dictionary.foundry
share/icradius/raddb/dictionary.livingston
share/icradius/raddb/dictionary.merit
share/icradius/raddb/dictionary.redback
share/icradius/raddb/dictionary.shiva
share/icradius/raddb/dictionary.tunnel
share/icradius/raddb/dictionary.usr
share/icradius/raddb/dictionary.versanet
share/icradius/raddb/huntgroups
share/icradius/scripts/images/l.gif
share/icradius/scripts/images/lc.gif
share/icradius/scripts/images/r.gif
share/icradius/scripts/images/rc.gif
share/icradius/scripts/images/rr.gif
share/icradius/scripts/radwho
share/icradius/scripts/radlast
share/icradius/scripts/radius.db
share/icradius/scripts/usonlineconv.pl
share/icradius/scripts/userimport.pl
share/icradius/scripts/userexport.pl
share/icradius/scripts/syncaccounting.pl
share/icradius/scripts/radiusfixup.pl
share/icradius/scripts/dictimport.pl
share/icradius/scripts/acctsummarize.pl
share/icradius/scripts/acctimport.pl
share/icradius/scripts/README
share/icradius/scripts/acctexport.pl
@dirrm etc/raddb
%%PORTDOCS%%@dirrm share/doc/icradius
@dirrm share/icradius/scripts/images
@dirrm share/icradius/scripts
@dirrm share/icradius/raddb
@dirrm share/icradius
