# mail/pysieved/files/pysieved.in
#!/bin/sh
#
# PROVIDE: %%PORTNAME%%
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown
#
# Define these %%PORTNAME%%_* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#

. /etc/rc.subr

name="%%PORTNAME%%"
rcvar=`set_rcvar`

load_rc_config $name

: ${%%PORTNAME%%_enable="NO"}
: ${%%PORTNAME%%_pidfile="/var/run/%%PORTNAME%%.pid"}

pidfile="${%%PORTNAME%%_pidfile}"
command_interpreter="%%PYTHON_CMD%%"
command="%%PREFIX%%/lib/pysieved/pysieved.py"

run_rc_command "$1"
