# sysutils/paladin/pkg-plist
bin/paladin
@sample etc/paladin.conf.sample
etc/rc.d/paladin
