# security/racoon2/pkg-plist
sbin/iked
%%KINK%%sbin/kinkd
sbin/spmd
sbin/spmdctl
sbin/pskgen
etc/spmd.pwd
etc/default.conf.sample
etc/local-test.conf.sample
etc/racoon2.conf.sample
etc/transport_ike.conf.sample
etc/transport_kink.conf.sample
etc/tunnel_ike.conf.sample
etc/tunnel_ike_natt.conf.sample
etc/tunnel_kink.conf.sample
etc/vals.conf.sample
etc/rc.d/iked.sh
etc/rc.d/spmd.sh
%%KINK%%etc/rc.d/kinkd.sh
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT.jp
%%PORTDOCS%%%%DOCSDIR%%/README
