# net/bird/files/bird6.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: bird6 dynamicrouting
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# bird6_enable (bool):   Set to NO by default.
#               Set it to YES to enable bird6.
# bird6_config (path):   Set to %%PREFIX%%/etc/bird6.conf
#               by default.
#

. /etc/rc.subr

name="bird6"
rcvar=bird6_enable

command=%%PREFIX%%/sbin/${name}

load_rc_config $name

extra_commands="reload"

: ${bird6_enable="NO"}
: ${bird6_config="%%PREFIX%%/etc/bird6.conf"}
: ${bird6_group="birdvty"}

command_args="-P /var/run/${name}.pid -c $bird6_config -g $bird6_group"

run_rc_command "$1"
