# japanese/sj3/pkg-plist
bin/sj3
bin/sj3dic
bin/sj3mkdic
bin/sj3serv
bin/sj3stat
include/sj3lib.h
lib/libsj3lib.a
lib/sj3/Nihongo.ps.Z
lib/sj3/cvtkey.news
lib/sj3/cvtkey.nwp511
lib/sj3/cvtkey.nwp517
lib/sj3/cvtkey.xterm
lib/sj3/dict/sj3main.dic
lib/sj3/hinsi.ms
lib/sj3/makedict.ms
lib/sj3/serverrc
lib/sj3/sj3lib.ms
lib/sj3/sj3serv.ms
lib/sj3/sjrc
lib/sj3/sjrk
etc/rc.d/sj3.sh
@dirrm lib/sj3/dict/user
@dirrm lib/sj3/dict
@dirrm lib/sj3
