# security/samhain/pkg-plist
@comment $FreeBSD$
%%YULE%%etc/yulerc.sample
%%YULE%%sbin/yule_setpwd
%%YULE%%sbin/yule
%%YULE%%etc/rc.d/yule.sh.sample
%%YULE%%@unexec rmdir /var/lib/yule 2>/dev/null || true
%%SAMHAIN%%sbin/samhain
%%SAMHAIN%%etc/rc.d/samhain.sh.sample
%%SAMHAIN%%etc/samhainrc.sample
%%SAMHAIN%%@unexec rmdir /var/lib/samhain 2>/dev/null || true
%%SETPWD%%sbin/samhain_setpwd
%%PORTDOCS%%%%DOCSDIR%%/MANUAL-2_0.ps
%%PORTDOCS%%%%DOCSDIR%%/MANUAL-2_0.html.tar
%%PORTDOCS%%%%DOCSDIR%%/HOWTO-client+server.html
%%PORTDOCS%%%%DOCSDIR%%/HOWTO-client+server-troubleshooting.html
%%PORTDOCS%%%%DOCSDIR%%/HOWTO-write-modules.html
%%PORTDOCS%%%%DOCSDIR%%/HOWTO-samhain+GnuPG.html
%%PORTDOCS%%%%DOCSDIR%%/FAQ.html
%%PORTDOCS%%%%DOCSDIR%%/README.UPGRADE
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/BUGS
%%PORTDOCS%%%%DOCSDIR%%/sh_mounts.txt
%%PORTDOCS%%%%DOCSDIR%%/sh_userfiles.txt
%%PORTDOCS%%@dirrm %%DOCSDIR%%
