# security/samhain/pkg-plist
@dir /var/lib
%%PORTDOCS%%%%DOCSDIR%%/BUGS
%%PORTDOCS%%%%DOCSDIR%%/FAQ.html
%%PORTDOCS%%%%DOCSDIR%%/HOWTO-client+server-troubleshooting.html
%%PORTDOCS%%%%DOCSDIR%%/HOWTO-client+server.html
%%PORTDOCS%%%%DOCSDIR%%/HOWTO-samhain+GnuPG.html
%%PORTDOCS%%%%DOCSDIR%%/HOWTO-write-modules.html
%%PORTDOCS%%%%DOCSDIR%%/MANUAL-2_4.pdf
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README.UPGRADE
%%PORTDOCS%%%%DOCSDIR%%/sh_mounts.txt
%%PORTDOCS%%%%DOCSDIR%%/sh_userfiles.txt
%%PORTDOCS%%@dir %%DOCSDIR%%
%%SAMHAIN%%@dir /var/lib/samhain
%%SAMHAIN%%@sample etc/samhainrc.sample
%%SAMHAIN%%etc/rc.d/samhain
%%SAMHAIN%%man/man5/samhainrc.5.gz
%%SAMHAIN%%man/man8/samhain.8.gz
%%SAMHAIN%%sbin/samhain
%%SETPWD%%sbin/samhain_setpwd
%%YULE%%@dir /var/lib
%%YULE%%@dir /var/lib/yule
%%YULE%%@dir /var/log/yule
%%YULE%%@sample etc/yulerc.sample
%%YULE%%@postunexec echo "To delete the yule user permanently, use 'pw userdel yule'"
%%YULE%%etc/rc.d/yule
%%YULE%%man/man5/yulerc.5.gz
%%YULE%%man/man8/yule.8.gz
%%YULE%%sbin/yule
%%YULE%%sbin/yule_setpwd
%%YULE%%sbin/yuleadmin-gpg.pl
%%YULE%%sbin/yuleadmin-sig.pl
%%YULE%%sbin/yulectl
