# databases/sqlrelay/pkg-plist
bin/backupschema
bin/fields
bin/query
bin/sqlr-cachemanager
bin/sqlrclient-config
bin/sqlrclientwrapper-config
%%IF_GTK%%bin/sqlr-config-gtk
%%IF_FREETDS%%bin/sqlr-connection-freetds
%%IF_FREETDS%%bin/sqlr-connection-freetds-debug
%%IF_MSQL%%bin/sqlr-connection-msql
%%IF_MSQL%%bin/sqlr-connection-msql-debug
%%IF_MYSQL%%bin/sqlr-connection-mysql
%%IF_MYSQL%%bin/sqlr-connection-mysql-debug
%%IF_ODBC%%bin/sqlr-connection-odbc
%%IF_ODBC%%bin/sqlr-connection-odbc-debug
%%IF_POSTGRESQL%%bin/sqlr-connection-postgresql
%%IF_POSTGRESQL%%bin/sqlr-connection-postgresql-debug
%%IF_SQLITE%%bin/sqlr-connection-sqlite
%%IF_SQLITE%%bin/sqlr-connection-sqlite-debug
%%IF_INTERBASE%%bin/sqlr-connection-interbase
%%IF_INTERBASE%%bin/sqlr-connection-interbase-debug
bin/sqlr-listener
bin/sqlr-listener-debug
bin/sqlr-scaler
bin/sqlr-start
bin/sqlr-stop
bin/sqlrsh
etc/rc.d/sqlrelay.sh-sample
etc/sqlrelay.conf.example
etc/sqlrelay.dtd
etc/sqlrelay-sample
include/sqlrelay/private/bindvar.h
include/sqlrelay/private/column.h
include/sqlrelay/private/row.h
include/sqlrelay/private/sqlrconnection.h
include/sqlrelay/private/sqlrcursor.h
include/sqlrelay/private/sqlrdefines.h
include/sqlrelay/private/sqlrincludes.h
include/sqlrelay/sqlrclient.h
include/sqlrelay/sqlrclientwrapper.h
@dirrm include/sqlrelay/private
@dirrm include/sqlrelay
lib/libmysql3sqlrelay.so.0
lib/libmysql3sqlrelay.so
lib/libmysql3sqlrelay.a
lib/libmysql40sqlrelay.so.0
lib/libmysql40sqlrelay.so
lib/libmysql40sqlrelay.a
lib/libmysql41sqlrelay.so.0
lib/libmysql41sqlrelay.so
lib/libmysql41sqlrelay.a
lib/libmysql50sqlrelay.so.0
lib/libmysql50sqlrelay.so
lib/libmysql50sqlrelay.a
lib/libsqlrclient.a
lib/libsqlrclient.so
lib/libsqlrclient.so.0
lib/libsqlrclientwrapper.a
lib/libsqlrclientwrapper.so
lib/libsqlrclientwrapper.so.0
lib/libpqsqlrelay.a
lib/libpqsqlrelay.so
lib/libpqsqlrelay.so.0
lib/libsqlrconnection.a
lib/libsqlrconnection.so
lib/libsqlrconnection.so.0
lib/libsqlrconnection_debug.a
lib/libsqlrconnection_debug.so
lib/libsqlrconnection_debug.so.0
lib/libsqlrutil.a
lib/libsqlrutil.so
lib/libsqlrutil.so.0
libdata/pkgconfig/sqlrelay-c.pc
libdata/pkgconfig/sqlrelay-c++.pc
%%DOCSDIR%%/api/c++/sqlrclient.h.html
%%DOCSDIR%%/api/c/sqlrclientwrapper.h.html
%%DOCSDIR%%/api/java/allclasses-frame.html
%%DOCSDIR%%/api/java/allclasses-noframe.html
%%DOCSDIR%%/api/java/com/firstworks/sqlrelay/SQLRConnection.html
%%DOCSDIR%%/api/java/com/firstworks/sqlrelay/SQLRCursor.html
%%DOCSDIR%%/api/java/com/firstworks/sqlrelay/package-frame.html
%%DOCSDIR%%/api/java/com/firstworks/sqlrelay/package-summary.html
%%DOCSDIR%%/api/java/com/firstworks/sqlrelay/package-tree.html
%%DOCSDIR%%/api/java/constant-values.html
%%DOCSDIR%%/api/java/deprecated-list.html
%%DOCSDIR%%/api/java/help-doc.html
%%DOCSDIR%%/api/java/index-all.html
%%DOCSDIR%%/api/java/index.html
%%DOCSDIR%%/api/java/overview-tree.html
%%DOCSDIR%%/api/java/package-list
%%DOCSDIR%%/api/java/packages.html
%%DOCSDIR%%/api/java/stylesheet.css
%%DOCSDIR%%/api/perl/Connection.html
%%DOCSDIR%%/api/perl/Cursor.html
%%DOCSDIR%%/api/perl/SQLRelay.html
%%DOCSDIR%%/api/perl/perl.html
%%DOCSDIR%%/api/php/doc.html
%%DOCSDIR%%/api/python/PySQLRClient.py.html
%%DOCSDIR%%/api/python/PySQLRDB.py.html
%%DOCSDIR%%/api/ruby/ruby.html
%%DOCSDIR%%/api/ruby/rubydbi.html
%%DOCSDIR%%/clientapis.html
%%DOCSDIR%%/configuring.html
%%DOCSDIR%%/connections.html
%%DOCSDIR%%/css/styles.css
%%DOCSDIR%%/dropin/mysql.html
%%DOCSDIR%%/dropin/postgresql.html
%%DOCSDIR%%/faq.html
%%DOCSDIR%%/gettingstarted/db2.html
%%DOCSDIR%%/gettingstarted/interbase.html
%%DOCSDIR%%/gettingstarted/lago.html
%%DOCSDIR%%/gettingstarted/msql.html
%%DOCSDIR%%/gettingstarted/mssqlserver.html
%%DOCSDIR%%/gettingstarted/mysql.html
%%DOCSDIR%%/gettingstarted/oracle.html
%%DOCSDIR%%/gettingstarted/oracle/linux/redhat/5/805.html
%%DOCSDIR%%/gettingstarted/oracle/linux/redhat/6/8i.html
%%DOCSDIR%%/gettingstarted/oracle/linux/redhat/7/8i.html
%%DOCSDIR%%/gettingstarted/oracle/linux/slackware/7/8i.html
%%DOCSDIR%%/gettingstarted/oracle/linux/slackware/8/8i.html
%%DOCSDIR%%/gettingstarted/oracle/linux/slackware/8/9i.html
%%DOCSDIR%%/gettingstarted/postgresql.html
%%DOCSDIR%%/gettingstarted/sqlite.html
%%DOCSDIR%%/gettingstarted/sybase.html
%%DOCSDIR%%/gettingstarted/sybase/11.0.3.3.html
%%DOCSDIR%%/gettingstarted/sybase/11.9.2.html
%%DOCSDIR%%/index.html
%%DOCSDIR%%/installing.html
%%DOCSDIR%%/introduction.html
%%DOCSDIR%%/introtowebdb.html
%%DOCSDIR%%/invisible.gif
%%DOCSDIR%%/oraclentier.html
%%DOCSDIR%%/programming/binds.html
%%DOCSDIR%%/programming/c++.html
%%DOCSDIR%%/programming/c.html
%%DOCSDIR%%/programming/java.html
%%DOCSDIR%%/programming/perl.html
%%DOCSDIR%%/programming/perldbi.html
%%DOCSDIR%%/programming/php.html
%%DOCSDIR%%/programming/python.html
%%DOCSDIR%%/programming/pythondb.html
%%DOCSDIR%%/programming/ruby.html
%%DOCSDIR%%/programming/rubydbi.html
%%DOCSDIR%%/programming/zope.html
%%DOCSDIR%%/programming/zope/instance.gif
%%DOCSDIR%%/programming/zope/management.gif
%%DOCSDIR%%/programming/zope/resultset.gif
%%DOCSDIR%%/programming/zope/sqlrelay.gif
%%DOCSDIR%%/programming/zope/test.gif
%%DOCSDIR%%/programming/zope/testzsqlmethod.gif
%%DOCSDIR%%/programming/zope/zsqlmethod.gif
%%DOCSDIR%%/programming/zope/zsqlmethodresultset.gif
%%DOCSDIR%%/running.html
%%DOCSDIR%%/tuning.html
%%DOCSDIR%%/windows.html
@unexec rmdir /var/sqlrelay/cache
@unexec rmdir /var/sqlrelay/debug
@unexec rmdir /var/sqlrelay/tmp/ipc
@unexec rmdir /var/sqlrelay/tmp/pids
@unexec rmdir /var/sqlrelay/tmp/sockets
@unexec rmdir /var/sqlrelay/tmp
@unexec rmdir /var/sqlrelay
@dirrm %%DOCSDIR%%/programming/zope
@dirrm %%DOCSDIR%%/programming
@dirrm %%DOCSDIR%%/gettingstarted/sybase
@dirrm %%DOCSDIR%%/gettingstarted/oracle/linux/slackware/8
@dirrm %%DOCSDIR%%/gettingstarted/oracle/linux/slackware/7
@dirrm %%DOCSDIR%%/gettingstarted/oracle/linux/slackware
@dirrm %%DOCSDIR%%/gettingstarted/oracle/linux/redhat/9
@dirrm %%DOCSDIR%%/gettingstarted/oracle/linux/redhat/8
@dirrm %%DOCSDIR%%/gettingstarted/oracle/linux/redhat/7
@dirrm %%DOCSDIR%%/gettingstarted/oracle/linux/redhat/6
@dirrm %%DOCSDIR%%/gettingstarted/oracle/linux/redhat/5
@dirrm %%DOCSDIR%%/gettingstarted/oracle/linux/redhat
@dirrm %%DOCSDIR%%/gettingstarted/oracle/linux
@dirrm %%DOCSDIR%%/gettingstarted/oracle
@dirrm %%DOCSDIR%%/gettingstarted
@dirrm %%DOCSDIR%%/dropin
@dirrm %%DOCSDIR%%/css
@dirrm %%DOCSDIR%%/api/ruby
@dirrm %%DOCSDIR%%/api/python
@dirrm %%DOCSDIR%%/api/php
@dirrm %%DOCSDIR%%/api/perl
@dirrm %%DOCSDIR%%/api/java/com/firstworks/sqlrelay
@dirrm %%DOCSDIR%%/api/java/com/firstworks
@dirrm %%DOCSDIR%%/api/java/com
@dirrm %%DOCSDIR%%/api/java
@dirrm %%DOCSDIR%%/api/c++
@dirrm %%DOCSDIR%%/api/c
@dirrm %%DOCSDIR%%/api
@dirrm %%DOCSDIR%%
