# databases/sqlrelay/pkg-plist
bin/backupschema
bin/sqlr-cachemanager
bin/sqlr-connection
bin/sqlr-export
bin/sqlr-fields
bin/sqlr-import
bin/sqlr-listener
bin/sqlr-pwdenc
bin/sqlr-query
bin/sqlr-scaler
bin/sqlr-start
bin/sqlr-status
bin/sqlr-stop
bin/sqlrclient-config
bin/sqlrclientwrapper-config
bin/sqlrserver-config
bin/sqlrsh
etc/rc.d/sqlrelay.sh-sample
etc/sqlrelay-sample
etc/sqlrelay.conf.example
etc/sqlrelay.dtd
etc/sqlrelay.xsd
include/sqlrelay/private/bindvar.h
include/sqlrelay/private/column.h
include/sqlrelay/private/dll.h
include/sqlrelay/private/row.h
include/sqlrelay/private/sqlrclientwrapper.h
include/sqlrelay/private/sqlrconnection.h
include/sqlrelay/private/sqlrcursor.h
include/sqlrelay/private/sqlrdefines.h
include/sqlrelay/private/sqlrincludes.h
include/sqlrelay/private/wrapperdll.h
include/sqlrelay/sqlrclient.h
include/sqlrelay/sqlrclientwrapper.h
include/sqlrelay/sqlrpwdenc.h
lib/libmysql3sqlrelay.so
lib/libmysql3sqlrelay.so.0
lib/libmysql40sqlrelay.so
lib/libmysql40sqlrelay.so.0
lib/libmysql41sqlrelay.so
lib/libmysql41sqlrelay.so.0
lib/libmysql50sqlrelay.so
lib/libmysql50sqlrelay.so.0
lib/libmysql51sqlrelay.so
lib/libmysql51sqlrelay.so.0
lib/libpqsqlrelay.so
lib/libpqsqlrelay.so.0
lib/libsqlrclient.a
lib/libsqlrclient.la
lib/libsqlrclient.so
lib/libsqlrclient.so.0
lib/libsqlrclientwrapper.a
lib/libsqlrclientwrapper.la
lib/libsqlrclientwrapper.so
lib/libsqlrclientwrapper.so.0
%%ODBC%%lib/libsqlrodbc.so
%%ODBC%%lib/libsqlrodbc.so.0
lib/libsqlrserver.a
lib/libsqlrserver.la
lib/libsqlrserver.so
lib/libsqlrserver.so.0
lib/libsqlrutil.a
lib/libsqlrutil.la
lib/libsqlrutil.so
lib/libsqlrutil.so.0
libdata/pkgconfig/sqlrelay-c++.pc
libdata/pkgconfig/sqlrelay-c.pc
%%FIREBIRD%%libexec/sqlrelay/sqlrconnection_firebird.so
%%FREETDS%%libexec/sqlrelay/sqlrconnection_freetds.so
%%MSQL%%libexec/sqlrelay/sqlrconnection_msql.so
%%MYSQL%%libexec/sqlrelay/sqlrconnection_mysql.so
%%ODBC%%libexec/sqlrelay/sqlrconnection_odbc.so
%%PGSQL%%libexec/sqlrelay/sqlrconnection_postgresql.so
%%SQLITE3%%libexec/sqlrelay/sqlrconnection_sqlite.so
libexec/sqlrelay/sqlrconnection_router.so
libexec/sqlrelay/sqlrlogger_custom_nw.so
libexec/sqlrelay/sqlrlogger_custom_sc.so
libexec/sqlrelay/sqlrlogger_debug.so
libexec/sqlrelay/sqlrlogger_slowqueries.so
libexec/sqlrelay/sqlrpwdenc_crypt.so
libexec/sqlrelay/sqlrpwdenc_md5.so
libexec/sqlrelay/sqlrpwdenc_rot.so
libexec/sqlrelay/sqlrquery_sqlrcmdcstat.so
libexec/sqlrelay/sqlrquery_sqlrcmdgstat.so
libexec/sqlrelay/sqltranslation_concat.so
libexec/sqlrelay/sqltranslation_doublequotestosinglequotes.so
libexec/sqlrelay/sqltranslation_forupdatemssqlserverize.so
libexec/sqlrelay/sqltranslation_informixtomssqlserverdate.so
libexec/sqlrelay/sqltranslation_informixtomssqlserverselectinto.so
libexec/sqlrelay/sqltranslation_informixtooracledate.so
libexec/sqlrelay/sqltranslation_locksmssqlserverize.so
libexec/sqlrelay/sqltranslation_locksnowaitbydefault.so
libexec/sqlrelay/sqltranslation_matchestolike.so
libexec/sqlrelay/sqltranslation_oracleize.so
libexec/sqlrelay/sqltranslation_serialtoautoincrement.so
libexec/sqlrelay/sqltranslation_serialtoidentity.so
libexec/sqlrelay/sqltranslation_show.so
libexec/sqlrelay/sqltranslation_temptablesaddmissingcolumns.so
libexec/sqlrelay/sqltranslation_temptableslocalize.so
libexec/sqlrelay/sqltranslation_temptablessybaseize.so
libexec/sqlrelay/sqltranslation_translatedatetimes.so
libexec/sqlrelay/sqltrigger_createtableautoincrementoracle.so
libexec/sqlrelay/sqltrigger_droplocalizedtemptables.so
libexec/sqlrelay/sqltrigger_droptableautoincrementoracle.so
man/man1/sqlr-fields.1.gz
man/man1/sqlr-query.1.gz
man/man1/sqlrsh.1.gz
man/man8/sqlr-cachemanager.8.gz
man/man8/sqlr-connection.8.gz
man/man8/sqlr-ipclean.8.gz
man/man8/sqlr-listener.8.gz
man/man8/sqlr-scaler.8.gz
man/man8/sqlr-start.8.gz
man/man8/sqlr-stop.8.gz
share/examples/sample-mssqlserver-mysql-errors.map
share/examples/sample-oracle-mysql-errors.map
@dirrm libexec/sqlrelay
@dirrm include/sqlrelay/private
@dirrm include/sqlrelay
@cwd /var
@dirrmtry sqlrelay/tmp/sockets
@dirrmtry sqlrelay/tmp/pids
@dirrmtry sqlrelay/tmp/ipc
@dirrmtry sqlrelay/tmp
@dirrmtry sqlrelay/debug
@dirrmtry sqlrelay/cache
@dirrmtry sqlrelay
