# databases/sqlrelay/pkg-plist
bin/sqlr-cachemanager
bin/sqlr-connection
bin/sqlr-export
bin/sqlr-import
bin/sqlr-listener
bin/sqlr-pwdenc
bin/sqlr-scaler
bin/sqlr-start
bin/sqlr-status
bin/sqlr-stop
bin/sqlrclient-config
bin/sqlrclientwrapper-config
bin/sqlrserver-config
bin/sqlrsh
etc/rc.d/sqlrcachemanager
etc/rc.d/sqlrelay
@sample etc/sqlrelay/conf.sample
etc/sqlrelay/sqlrelay.xsd
include/sqlrelay/private/sqlrclientincludes.h
include/sqlrelay/private/sqlrclientwrapper.h
include/sqlrelay/private/sqlrclientwrapperincludes.h
include/sqlrelay/private/sqlrconnection.h
include/sqlrelay/private/sqlrcursor.h
include/sqlrelay/private/sqlrlistener.h
include/sqlrelay/private/sqlrserverconnection.h
include/sqlrelay/private/sqlrservercontroller.h
include/sqlrelay/private/sqlrservercursor.h
include/sqlrelay/private/sqlrserverincludes.h
include/sqlrelay/private/sqlrutilincludes.h
include/sqlrelay/private/sqlrshmdata.h
include/sqlrelay/sqlrclient.h
include/sqlrelay/sqlrclientwrapper.h
include/sqlrelay/sqlrserver.h
include/sqlrelay/sqlrutil.h
lib/libmysql3sqlrelay.so
lib/libmysql3sqlrelay.so.2
lib/libmysql3sqlrelay.so.2.0.0
lib/libmysql40sqlrelay.so
lib/libmysql40sqlrelay.so.2
lib/libmysql40sqlrelay.so.2.0.0
lib/libmysql41sqlrelay.so
lib/libmysql41sqlrelay.so.2
lib/libmysql41sqlrelay.so.2.0.0
lib/libmysql50sqlrelay.so
lib/libmysql50sqlrelay.so.2
lib/libmysql50sqlrelay.so.2.0.0
lib/libmysql51sqlrelay.so
lib/libmysql51sqlrelay.so.2
lib/libmysql51sqlrelay.so.2.0.0
lib/libpqsqlrelay.so
lib/libpqsqlrelay.so.2
lib/libpqsqlrelay.so.2.0.0
lib/libsqlrclient.a
lib/libsqlrclient.so
lib/libsqlrclient.so.2
lib/libsqlrclient.so.2.0.0
lib/libsqlrclientwrapper.a
lib/libsqlrclientwrapper.so
lib/libsqlrclientwrapper.so.2
lib/libsqlrclientwrapper.so.2.0.0
%%ODBC%%lib/libsqlrodbc.so
%%ODBC%%lib/libsqlrodbc.so.2
%%ODBC%%lib/libsqlrodbc.so.2.0.0
lib/libsqlrserver.a
lib/libsqlrserver.so
lib/libsqlrserver.so.2
lib/libsqlrserver.so.2.0.0
lib/libsqlrutil.a
lib/libsqlrutil.so
lib/libsqlrutil.so.2
lib/libsqlrutil.so.2.0.0
libdata/pkgconfig/sqlrelay-c++.pc
libdata/pkgconfig/sqlrelay-c.pc
libexec/sqlrelay/sqlrauth_database.so
libexec/sqlrelay/sqlrauth_krb_userlist.so
libexec/sqlrelay/sqlrauth_proxied.so
libexec/sqlrelay/sqlrauth_tls_userlist.so
libexec/sqlrelay/sqlrauth_userlist.so
libexec/sqlrelay/sqlrconfig_xmldom.so
libexec/sqlrelay/sqlrfilter_patterns.so
libexec/sqlrelay/sqlrfilter_regex.so
libexec/sqlrelay/sqlrfilter_string.so
libexec/sqlrelay/sqlrresultsettranslation_reformatdatetime.so
libexec/sqlrelay/sqlrtranslation_normalize.so
%%FIREBIRD%%libexec/sqlrelay/sqlrconnection_firebird.so
%%FREETDS%%libexec/sqlrelay/sqlrconnection_freetds.so
%%MYSQL%%libexec/sqlrelay/sqlrconnection_mysql.so
%%ODBC%%libexec/sqlrelay/sqlrconnection_odbc.so
%%PGSQL%%libexec/sqlrelay/sqlrconnection_postgresql.so
%%SQLITE3%%libexec/sqlrelay/sqlrconnection_sqlite.so
libexec/sqlrelay/sqlrauth_sqlrelay.so
libexec/sqlrelay/sqlrconnection_router.so
libexec/sqlrelay/sqlrlogger_custom_nw.so
libexec/sqlrelay/sqlrlogger_custom_sc.so
libexec/sqlrelay/sqlrlogger_debug.so
libexec/sqlrelay/sqlrlogger_slowqueries.so
libexec/sqlrelay/sqlrparser_default.so
libexec/sqlrelay/sqlrprotocol_sqlrclient.so
libexec/sqlrelay/sqlrpwdenc_crypt.so
libexec/sqlrelay/sqlrpwdenc_md5.so
libexec/sqlrelay/sqlrpwdenc_rot.so
libexec/sqlrelay/sqlrquery_sqlrcmdcstat.so
libexec/sqlrelay/sqlrquery_sqlrcmdgstat.so
man/man1/sqlr-export.1.gz
man/man1/sqlr-import.1.gz
man/man1/sqlrsh.1.gz
man/man8/sqlr-cachemanager.8.gz
man/man8/sqlr-connection.8.gz
man/man8/sqlr-listener.8.gz
man/man8/sqlr-pwdenc.8.gz
man/man8/sqlr-scaler.8.gz
man/man8/sqlr-start.8.gz
man/man8/sqlr-status.8.gz
man/man8/sqlr-stop.8.gz
%%DATADIR%%/examples/sqlrelay.conf
%%DATADIR%%/examples/sample-mssqlserver-mysql-errors.map
%%DATADIR%%/examples/sample-oracle-mysql-errors.map
@dir %%ETCDIR%%
@dir /var/sqlrelay/tmp/sockets
@dir /var/sqlrelay/tmp/pids
@dir /var/sqlrelay/tmp/ipc
@dir /var/sqlrelay/tmp
@dir /var/sqlrelay/log
@dir /var/sqlrelay/debug
@dir /var/sqlrelay/cache
@dir /var/sqlrelay
