# databases/sqlrelay/pkg-plist
%%DATADIR%%/examples/sqlrelay.conf
%%FIREBIRD%%libexec/sqlrelay/sqlrconnection_firebird.so
%%FREETDS%%libexec/sqlrelay/sqlrconnection_freetds.so
%%MYSQL%%libexec/sqlrelay/sqlrconnection_mysql.so
%%ODBC%%lib/libsqlrodbc.so
%%ODBC%%lib/libsqlrodbc.so.5
%%ODBC%%lib/libsqlrodbc.so.5.0.0
%%ODBC%%libexec/sqlrelay/sqlrconnection_odbc.so
%%PGSQL%%libexec/sqlrelay/sqlrconnection_postgresql.so
%%SQLITE3%%libexec/sqlrelay/sqlrconnection_sqlite.so
@dir %%ETCDIR%%
@dir /var/cache/sqlrelay
@dir /var/log/sqlrelay
@dir /var/log/sqlrelay/debug
@dir /var/run/sqlrelay
@sample etc/sqlrelay/conf.sample
bin/sqlr-cachemanager
bin/sqlr-connection
bin/sqlr-export
bin/sqlr-import
bin/sqlr-listener
bin/sqlr-pwdenc
bin/sqlr-scaler
bin/sqlr-start
bin/sqlr-status
bin/sqlr-stop
bin/sqlrclient-config
bin/sqlrclientwrapper-config
bin/sqlrserver-config
bin/sqlrsh
etc/rc.d/sqlrcachemanager
etc/rc.d/sqlrelay
etc/sqlrelay/sqlrelay.xsd
include/sqlrelay/private/sqlrauth.h
include/sqlrelay/private/sqlrauths.h
include/sqlrelay/private/sqlrbindvariabletranslation.h
include/sqlrelay/private/sqlrbindvariabletranslations.h
include/sqlrelay/private/sqlrclientincludes.h
include/sqlrelay/private/sqlrclientwrapper.h
include/sqlrelay/private/sqlrclientwrapperincludes.h
include/sqlrelay/private/sqlrconnection.h
include/sqlrelay/private/sqlrcursor.h
include/sqlrelay/private/sqlrdirective.h
include/sqlrelay/private/sqlrdirectives.h
include/sqlrelay/private/sqlrfilter.h
include/sqlrelay/private/sqlrfilters.h
include/sqlrelay/private/sqlrgsscredentials.h
include/sqlrelay/private/sqlrlistener.h
include/sqlrelay/private/sqlrlogger.h
include/sqlrelay/private/sqlrloggers.h
include/sqlrelay/private/sqlrmoduledata.h
include/sqlrelay/private/sqlrmoduledata_tag.h
include/sqlrelay/private/sqlrmoduledatas.h
include/sqlrelay/private/sqlrmysqlcredentials.h
include/sqlrelay/private/sqlrnotification.h
include/sqlrelay/private/sqlrnotifications.h
include/sqlrelay/private/sqlrparser.h
include/sqlrelay/private/sqlrpostgresqlcredentials.h
include/sqlrelay/private/sqlrprotocol.h
include/sqlrelay/private/sqlrprotocols.h
include/sqlrelay/private/sqlrpwdenc.h
include/sqlrelay/private/sqlrpwdencs.h
include/sqlrelay/private/sqlrqueries.h
include/sqlrelay/private/sqlrquery.h
include/sqlrelay/private/sqlrquerycursor.h
include/sqlrelay/private/sqlrresultsetheadertranslation.h
include/sqlrelay/private/sqlrresultsetheadertranslations.h
include/sqlrelay/private/sqlrresultsetrowblocktranslation.h
include/sqlrelay/private/sqlrresultsetrowblocktranslations.h
include/sqlrelay/private/sqlrresultsetrowtranslation.h
include/sqlrelay/private/sqlrresultsetrowtranslations.h
include/sqlrelay/private/sqlrresultsettranslation.h
include/sqlrelay/private/sqlrresultsettranslations.h
include/sqlrelay/private/sqlrrouter.h
include/sqlrelay/private/sqlrrouters.h
include/sqlrelay/private/sqlrschedule.h
include/sqlrelay/private/sqlrschedulerule.h
include/sqlrelay/private/sqlrschedules.h
include/sqlrelay/private/sqlrserverconnection.h
include/sqlrelay/private/sqlrservercontroller.h
include/sqlrelay/private/sqlrservercursor.h
include/sqlrelay/private/sqlrserverincludes.h
include/sqlrelay/private/sqlrshm.h
include/sqlrelay/private/sqlrtlscredentials.h
include/sqlrelay/private/sqlrtranslation.h
include/sqlrelay/private/sqlrtranslations.h
include/sqlrelay/private/sqlrtrigger.h
include/sqlrelay/private/sqlrtriggers.h
include/sqlrelay/private/sqlruserpasswordcredentials.h
include/sqlrelay/private/sqlrutilincludes.h
include/sqlrelay/sqlrclient.h
include/sqlrelay/sqlrclientwrapper.h
include/sqlrelay/sqlrserver.h
include/sqlrelay/sqlrutil.h
lib/libsqlrclient.a
lib/libsqlrclient.so
lib/libsqlrclient.so.5
lib/libsqlrclient.so.5.0.0
lib/libsqlrclientwrapper.a
lib/libsqlrclientwrapper.so
lib/libsqlrclientwrapper.so.5
lib/libsqlrclientwrapper.so.5.0.0
lib/libsqlrserver.a
lib/libsqlrserver.so
lib/libsqlrserver.so.10
lib/libsqlrserver.so.10.0.0
lib/libsqlrutil.a
lib/libsqlrutil.so
lib/libsqlrutil.so.10
lib/libsqlrutil.so.10.0.0
libdata/pkgconfig/sqlrelay-c++.pc
libdata/pkgconfig/sqlrelay-c.pc
libexec/sqlrelay/sqlrauth_database.so
libexec/sqlrelay/sqlrauth_mysql_database.so
libexec/sqlrelay/sqlrauth_mysql_userlist.so
libexec/sqlrelay/sqlrauth_postgresql_userlist.so
libexec/sqlrelay/sqlrauth_proxied.so
libexec/sqlrelay/sqlrauth_sqlrelay.so
libexec/sqlrelay/sqlrauth_userlist.so
libexec/sqlrelay/sqlrbindvariabletranslation_test.so
libexec/sqlrelay/sqlrconfig_xmldom.so
libexec/sqlrelay/sqlrconnection_router.so
libexec/sqlrelay/sqlrdirective_crash.so
libexec/sqlrelay/sqlrdirective_custom_wf.so
libexec/sqlrelay/sqlrfilter_patterns.so
libexec/sqlrelay/sqlrfilter_regex.so
libexec/sqlrelay/sqlrfilter_string.so
libexec/sqlrelay/sqlrfilter_tag.so
libexec/sqlrelay/sqlrlogger_custom_nw.so
libexec/sqlrelay/sqlrlogger_custom_sc.so
libexec/sqlrelay/sqlrlogger_debug.so
libexec/sqlrelay/sqlrlogger_slowqueries.so
libexec/sqlrelay/sqlrlogger_sql.so
libexec/sqlrelay/sqlrlogger_stalecursors.so
libexec/sqlrelay/sqlrmoduledata_tag.so
libexec/sqlrelay/sqlrnotification_events.so
libexec/sqlrelay/sqlrparser_default.so
libexec/sqlrelay/sqlrprotocol_mysql.so
libexec/sqlrelay/sqlrprotocol_postgresql.so
libexec/sqlrelay/sqlrprotocol_sqlrclient.so
libexec/sqlrelay/sqlrpwdenc_crypt.so
libexec/sqlrelay/sqlrpwdenc_md5.so
libexec/sqlrelay/sqlrpwdenc_rot.so
libexec/sqlrelay/sqlrpwdenc_sha1.so
libexec/sqlrelay/sqlrpwdenc_sha256.so
libexec/sqlrelay/sqlrquery_sqlrcmdcstat.so
libexec/sqlrelay/sqlrquery_sqlrcmdgstat.so
libexec/sqlrelay/sqlrresultsettranslation_reformatdatetime.so
libexec/sqlrelay/sqlrrouter_clientinfolist.so
libexec/sqlrelay/sqlrrouter_clientiplist.so
libexec/sqlrelay/sqlrrouter_regex.so
libexec/sqlrelay/sqlrrouter_usedatabase.so
libexec/sqlrelay/sqlrrouter_userlist.so
libexec/sqlrelay/sqlrschedule_cron_userlist.so
libexec/sqlrelay/sqlrtranslation_normalize.so
libexec/sqlrelay/sqlrtranslation_patterns.so
libexec/sqlrelay/sqlrtrigger_replay.so
man/man1/sqlr-export.1.gz
man/man1/sqlr-import.1.gz
man/man1/sqlrsh.1.gz
man/man8/sqlr-cachemanager.8.gz
man/man8/sqlr-connection.8.gz
man/man8/sqlr-listener.8.gz
man/man8/sqlr-pwdenc.8.gz
man/man8/sqlr-scaler.8.gz
man/man8/sqlr-start.8.gz
man/man8/sqlr-status.8.gz
man/man8/sqlr-stop.8.gz
