# net/gatekeeper/pkg-plist
@unexec %D/etc/rc.d/gnugk.sh forcestop 2>/dev/null || true
bin/addpasswd
bin/gnugk
etc/gnugk.ini.default
%%EXAMPLESDIR%%/%%PORTNAME%%/child.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/clirw.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/complete.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/fileacct.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/gnugk.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/ipauth.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/mysql.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/neighbor.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/proxy.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/radacct.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/radaliasauth.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/radauth.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/routed.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/rrqauth.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/sql.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/status.ini
%%EXAMPLESDIR%%/%%PORTNAME%%/vqueue.ini
@dirrm %%EXAMPLESDIR%%/%%PORTNAME%%
@dirrmtry %%EXAMPLESDIR%%
%%PORTDOCS%%%%DOCSDIR%%/%%PORTNAME%%/keepalive.txt
%%PORTDOCS%%%%DOCSDIR%%/%%PORTNAME%%/known_bugs.txt
%%PORTDOCS%%%%DOCSDIR%%/%%PORTNAME%%/performance.txt
%%PORTDOCS%%%%DOCSDIR%%/%%PORTNAME%%/radauth.txt
%%PORTDOCS%%%%DOCSDIR%%/%%PORTNAME%%/trace_levels.txt
%%PORTDOCS%%@dirrm %%DOCSDIR%%/%%PORTNAME%%
%%PORTDOCS%%@dirrmtry %%DOCSDIR%%
