# net/gatekeeper/pkg-plist
@unexec %D/etc/rc.d/gnugk.sh forcestop 2>/dev/null || true
bin/gnugk
bin/addpasswd
etc/rc.d/gnugk.sh
etc/gnugk.ini.default
share/examples/gatekeeper/child.ini
share/examples/gatekeeper/clirw.ini
share/examples/gatekeeper/complete.ini
share/examples/gatekeeper/fileacct.ini
share/examples/gatekeeper/mysql.ini
share/examples/gatekeeper/neighbor.ini
share/examples/gatekeeper/proxy.ini
share/examples/gatekeeper/radacct.ini
share/examples/gatekeeper/radaliasauth.ini
share/examples/gatekeeper/radauth.ini
share/examples/gatekeeper/routed.ini
share/examples/gatekeeper/rrqauth.ini
share/examples/gatekeeper/sql.ini
share/examples/gatekeeper/status.ini
share/examples/gatekeeper/vqueue.ini
@dirrm share/examples/gatekeeper
