# comms/scud/pkg-plist
sbin/scud
etc/rc.d/scud.sh.sample
etc/scud.conf.sample
