# net/isc-dhcp2/pkg-plist
sbin/dhcpd
sbin/dhcrelay
etc/rc.d/%%STARTUP_SCRIPT%%
@exec touch /var/db/dhcpd.leases
@exec touch %D/etc/dhclient.conf
