# ftp/smbftpd/files/openrc-smbftpd.in
#!/sbin/openrc-run

name=smbftpd
command=%%PREFIX%%/sbin/smbftpd
smbftpd_config=${smbftpd_config:-"%%PREFIX%%/etc/smbftpd/smbftpd.conf"}
required_files=${smbftpd_config}
pidfile=/var/run/smbftpd.pid
command_args="-D -s ${smbftpd_config}"

depend() {
	use net
	keyword -shutdown
}
