# net/rsmb/files/rsmbd.in
#!/bin/sh

# PROVIDE: rsmbd
# REQUIRE: NETWORKING
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# rsmbd_enable (bool):	Set to NO by default.
#			Set it to YES to enable rsmbd.
# rsmbd_config (path):	Unset by default.
#

. /etc/rc.subr

name="rsmbd"
rcvar=`set_rcvar`
eval "${name}_installdir=\${${name}_installdir:-'%%PREFIX%%/libexec/rsmb'}"
eval "command=\${${name}_installdir}/broker"

start_cmd="${name}_start"

load_rc_config $name

eval "${rcvar}=\${${rcvar}:-'NO'}"

rsmbd_start() {
	echo "Starting ${name}."
	cd %%PREFIX%%/libexec/rsmb && %%PREFIX%%/libexec/rsmb/broker ${rsmbd_config} > /dev/null 2>&1 &
}

run_rc_command "$1"

