# databases/postgresql-repmgr/files/repmgrd.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: repmgrd
# REQUIRE: postgresql
# KEYWORD: shutdown
#
# Add these lines to /etc/rc.conf.local or /etc/rc.conf to enable this
# service:
#
# repmgrd_enable (bool):   Set to YES to enable repmgrd.  Default: NO
#
# repmgrd_config (path):   Default: %%PREFIX%%/etc/repmgr.conf
#

. /etc/rc.subr

name=repmgrd
rcvar=repmgrd_enable

load_rc_config $name

: ${repmgrd_enable:="NO"}
: ${repmgrd_user:="pgsql"}
: ${repmgrd_config="%%PREFIX%%/etc/repmgr.conf"}

command=%%PREFIX%%/bin/${name}
pidfile=/var/run/${name}.pid

command_args="-m -d -p $pidfile -f $repmgrd_config"

run_rc_command "$1"

