# net/coturn/files/turnserver-service.in
#!/bin/sh
#
# $FreeBSD$
#

#
# PROVIDE: turnserver
# REQUIRE: DAEMON
# BEFORE: mail
# KEYWORD: shutdown
#

. /etc/rc.subr

name="turnserver"
rcvar="${name}_enable"
command="%%PREFIX%%/bin/${name}"
pidfile="/var/run/prosody/${name}.pid"
required_files="%%PREFIX%%/etc/${name}.conf"

load_rc_config ${name}

: ${turnserver_enable="NO"}
: ${turnserver_flags="-o"}
: ${turnserver_configfile:=%%PREFIX%%/etc/$name.conf}

command_args="-c ${turnserver_configfile} ${turnserver_flags}"

run_rc_command "$1"
