# net/isc-dhcp31-server/pkg-plist
@comment $FreeBSD$
bin/omshell
etc/dhcpd.conf.sample
etc/rc.d/isc-dhcpd.sh
sbin/dhcpd
