# net-im/jabber-pymsn/files/jabber-pymsn-transport.sh.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: jabber_pymsn
# REQUIRE: DAEMON %%JABBER_REQUIRE%%
# KEYWORD: shutdown
#
# Define these jabber_pymsn_* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#

. /etc/rc.subr

name="jabber_pymsn"
rcvar=jabber_pymsn_enable

load_rc_config $name

: ${jabber_pymsn_enable="NO"}
: ${jabber_pymsn_dir="%%PREFIX%%/lib/jabber/pymsn"}
: ${jabber_pymsn_piddir="%%JABBER_PIDDIR%%"}
: ${jabber_pymsn_user="%%JABBER_USER%%"}

pidfile="${jabber_pymsn_piddir}/PyMSNt.pid"
command_interpreter="%%PYTHON_CMD%%"
command_args="-p ${pidfile}"
command="${jabber_pymsn_dir}/PyMSNt.py"

run_rc_command "$1"
