# net-p2p/opennap/pkg-plist
sbin/opennap
sbin/metaserver
sbin/setup
share/opennap/users
share/opennap/examples/sample.block
share/opennap/examples/sample.channels
share/opennap/examples/sample.config
share/opennap/examples/sample.filter
share/opennap/examples/sample.motd
share/opennap/examples/sample.servers
share/opennap/examples/sample.users
%%PORTDOCS%%share/doc/opennap/README
%%PORTDOCS%%share/doc/opennap/napster.txt
@dirrm share/opennap/examples
@dirrm share/opennap
%%PORTDOCS%%@dirrm share/doc/opennap
etc/rc.d/opennap.sh
