# textproc/cocoon/pkg-plist
cocoon/CREDITS.txt
cocoon/INSTALL.txt
cocoon/LICENSE.txt
cocoon/NOTICE.txt
cocoon/README.txt
cocoon/legal/altrmi-client-impl-0.9.2.jar.license.txt
cocoon/legal/altrmi-client-interfaces-0.9.2.jar.license.txt
cocoon/legal/altrmi-common-0.9.2.jar.license.txt
cocoon/legal/altrmi-registry-0.9.2.jar.license.txt
cocoon/legal/altrmi-server-impl-0.9.2.jar.license.txt
cocoon/legal/altrmi-server-interfacesl-0.9.2.jar.license.txt
cocoon/legal/ant-contrib-0.6.jar.license.txt
cocoon/legal/ant-junit.jar.license.txt
cocoon/legal/ant-launcher.jar.license.txt
cocoon/legal/ant-trax.jar.license.txt
cocoon/legal/ant.jar.license.txt
cocoon/legal/antlr-2.7.3.jar.license.txt
cocoon/legal/asm-1.4.2.jar.license.txt
cocoon/legal/avalon-framework-api-4.1.5.jar.license.txt
cocoon/legal/avalon-framework-impl-4.1.5.jar.license.txt
cocoon/legal/axis-1.1.jar.license.txt
cocoon/legal/axis-jaxrpc-1.1.jar.license.txt
cocoon/legal/axis-saaj-1.1.jar.license.txt
cocoon/legal/batik-all-1.5.1.jar.license.txt
cocoon/legal/bsf-2.3.0.jar.license.txt
cocoon/legal/castor-0.9.5.3-xml.jar.license.txt
cocoon/legal/chaperon-20040205.jar.license.txt
cocoon/legal/commons-beanutils-1.6.1.jar.license.txt
cocoon/legal/commons-betwixt-20030910.jar.license.txt
cocoon/legal/commons-cli-1.0.jar.license.txt
cocoon/legal/commons-collections-3.0.jar.license.txt
cocoon/legal/commons-dbcp-1.1.jar.license.txt
cocoon/legal/commons-digester-1.3.jar.license.txt
cocoon/legal/commons-discovery-0.2.jar.license.txt
cocoon/legal/commons-httpclient-2.0-final.jar.license.txt
cocoon/legal/commons-jexl-1.0-beta-1-20040113.jar.license.txt
cocoon/legal/commons-jxpath-20030909.jar.license.txt
cocoon/legal/commons-lang-2.0.jar.license.txt
cocoon/legal/commons-logging-1.0.3.jar.license.txt
cocoon/legal/commons-pool-1.1.jar.license.txt
cocoon/legal/db-ojb-1.0.rc6.jar.license.txt
cocoon/legal/deli-0.9.8.jar.license.txt
cocoon/legal/ehcache-0.8.jar.license.txt
cocoon/legal/excalibur-component-1.2.jar.license.txt
cocoon/legal/excalibur-datasource-1.1.1.jar.license.txt
cocoon/legal/excalibur-event-api-1.1.jar.license.txt
cocoon/legal/excalibur-event-impl-1.1.jar.license.txt
cocoon/legal/excalibur-i18n-1.1.jar.license.txt
cocoon/legal/excalibur-instrument-1.0.jar.license.txt
cocoon/legal/excalibur-instrument-client-2003-03-31.jar.license.txt
cocoon/legal/excalibur-instrument-manager-1.0.jar.license.txt
cocoon/legal/excalibur-instrument-manager-interfaces-1.0.jar.license.txt
cocoon/legal/excalibur-io-1.1.jar.license.txt
cocoon/legal/excalibur-logger-1.1.jar.license.txt
cocoon/legal/excalibur-monitor-1.0.2.jar.license.txt
cocoon/legal/excalibur-naming-1.0.jar.license.txt
cocoon/legal/excalibur-pool-1.2.jar.license.txt
cocoon/legal/excalibur-sourceresolve-1.1.jar.license.txt
cocoon/legal/excalibur-store-1.0.jar.license.txt
cocoon/legal/excalibur-testcase-1.2.jar.license.txt
cocoon/legal/excalibur-util-1.0.jar.license.txt
cocoon/legal/excalibur-xmlutil-1.0.jar.license.txt
cocoon/legal/fop-0.20.5.jar.license.txt
cocoon/legal/geronimo-spec-jms-DEV-20031120.jar.license.txt
cocoon/legal/geronimo-spec-jta-DEV-20040202.jar.license.txt
cocoon/legal/groovy-20040415.114632.jar.license.txt
cocoon/legal/hsqldb-1.7.1.jar.license.txt
cocoon/legal/htmlarea.license.txt
cocoon/legal/icu4j.jar.license.txt
cocoon/legal/itext-1.02b.jar.license.txt
cocoon/legal/itext-xml-1.02.jar.license.txt
cocoon/legal/jakarta-bcel-20040329.jar.license.txt
cocoon/legal/jakarta-oro-2.0.8.jar.license.txt
cocoon/legal/jakarta-regexp-1.3.jar.license.txt
cocoon/legal/jakarta-slide-webdavlib-2.0rc1.jar.license.txt
cocoon/legal/jcs-1.0-dev.20040516.jar.license.txt
cocoon/legal/jdom-b9.jar.license.txt
cocoon/legal/jdtcore-2.1.3.jar.license.txt
cocoon/legal/jena-2.0.0.jar.license.txt
cocoon/legal/jetty-4.2.19.jar.license.html
cocoon/legal/jfor-0.7.1.jar.license.txt
cocoon/legal/jing-20030619.jar.license.txt
cocoon/legal/jisp-2.5.1.jar.license.txt
cocoon/legal/joost-20040330.jar.license.txt
cocoon/legal/jstyle.jar.license.txt
cocoon/legal/jtidy-04aug2000r7-dev.jar.license.txt
cocoon/legal/junit.jar.license.txt
cocoon/legal/log4j-1.2.8.jar.license.txt
cocoon/legal/logkit-1.2.2.jar.license.txt
cocoon/legal/lucene-1.3-final.jar.license.txt
cocoon/legal/midi-license.txt
cocoon/legal/pizza-1.1.jar.license.txt
cocoon/legal/pluto-20040310.jar.license.txt
cocoon/legal/poi-2.5-final-20040302.jar.license.txt
cocoon/legal/portlet-api-20040310.jar.license.txt
cocoon/legal/qdox-1.3.jar.license.txt
cocoon/legal/quartz-1.3.2.jar.license.txt
cocoon/legal/rhino1.5r4-continuations-20040228.jar.license.txt
cocoon/legal/servlet-2.3.jar.license.txt
cocoon/legal/servlet_2_2.jar.license.txt
cocoon/legal/slide-kernel-2.0.jar.license.txt
cocoon/legal/slide-stores-2.0.jar.license.txt
cocoon/legal/slide-webdavservlet-2.0.jar.license.txt
cocoon/legal/spark.fpl.license.txt
cocoon/legal/spark.opl.license.txt
cocoon/legal/util.concurrent-1.3.3.jar.license.txt
cocoon/legal/velocity-1.4.jar.license.jtxt
cocoon/legal/wsdl4j-1.4.jar.license.txt
cocoon/legal/xalan-2.6.0.jar.license.txt
cocoon/legal/xercesImpl-2.6.2.jar.license.txt
cocoon/legal/xindice-1.1b4.jar.license.txt
cocoon/legal/xml-apis.jar.license.txt
cocoon/legal/xml-commons-resolver-1.1.jar.license.txt
cocoon/legal/xmldb-api-20030701.jar.license.txt
cocoon/legal/xmldb-common-20030701.jar.license.txt
cocoon/legal/xmldb-xupdate-20040205.jar.license.txt
cocoon/legal/xmlunit0.8.jar.license.txt
cocoon/legal/xreporter-expression-20030725.jar.license.txt
cocoon/tools/jetty/conf/admin.xml
cocoon/tools/jetty/conf/main.xml
cocoon/tools/jetty/conf/webdefaults.xml
cocoon/tools/jetty/lib/jetty-4.2.19.jar
cocoon/tools/jetty/lib/servlet-2.3.jar
cocoon/tools/loader/Loader$RepositoryClassLoader.class
cocoon/tools/loader/Loader.class
@dirrm cocoon/legal
@dirrm cocoon/tools/jetty/conf
@dirrm cocoon/tools/jetty/lib
@dirrm cocoon/tools/jetty
@dirrm cocoon/tools/loader
@dirrm cocoon/tools
sbin/cocoon.sh
sbin/cocoonctl
etc/rc.d/cocoon.sh
@dirrm cocoon
