# net-mgmt/net-snmp/pkg-plist
bin/encode_keychange
bin/fixproc
bin/ipf-mod.pl
bin/mib2c
bin/net-snmp-config
bin/snmpbulkget
bin/snmpbulkwalk
bin/snmpcheck
bin/snmpconf
bin/snmpdelta
bin/snmpdf
bin/snmpget
bin/snmpgetnext
bin/snmpinform
bin/snmpnetstat
bin/snmpset
bin/snmpstatus
bin/snmptable
bin/snmptest
bin/snmptranslate
bin/snmptrap
bin/snmpusm
bin/snmpvacm
bin/snmpwalk
bin/traptoemail
etc/rc.d/snmpd.sh
include/net-snmp/agent/agent_callbacks.h
include/net-snmp/agent/agent_handler.h
include/net-snmp/agent/agent_index.h
include/net-snmp/agent/agent_read_config.h
include/net-snmp/agent/agent_registry.h
include/net-snmp/agent/agent_trap.h
include/net-snmp/agent/all_helpers.h
include/net-snmp/agent/auto_nlist.h
include/net-snmp/agent/bulk_to_next.h
include/net-snmp/agent/debug_handler.h
include/net-snmp/agent/ds_agent.h
include/net-snmp/agent/instance.h
include/net-snmp/agent/mib_module_config.h
include/net-snmp/agent/mib_module_includes.h
include/net-snmp/agent/multiplexer.h
include/net-snmp/agent/net-snmp-agent-includes.h
include/net-snmp/agent/null.h
include/net-snmp/agent/old_api.h
include/net-snmp/agent/read_only.h
include/net-snmp/agent/scalar.h
include/net-snmp/agent/serialize.h
include/net-snmp/agent/set_helper.h
include/net-snmp/agent/snmp_agent.h
include/net-snmp/agent/snmp_vars.h
include/net-snmp/agent/table.h
include/net-snmp/agent/table_array.h
include/net-snmp/agent/table_data.h
include/net-snmp/agent/table_dataset.h
include/net-snmp/agent/table_iterator.h
include/net-snmp/agent/var_struct.h
include/net-snmp/agent/watcher.h
include/net-snmp/config_api.h
include/net-snmp/definitions.h
include/net-snmp/library/README
include/net-snmp/library/asn1.h
include/net-snmp/library/callback.h
include/net-snmp/library/check_varbind.h
include/net-snmp/library/cmu_compat.h
include/net-snmp/library/container.h
include/net-snmp/library/container_binary_array.h
include/net-snmp/library/data_list.h
include/net-snmp/library/default_store.h
include/net-snmp/library/factory.h
include/net-snmp/library/getopt.h
include/net-snmp/library/int64.h
include/net-snmp/library/keytools.h
include/net-snmp/library/lcd_time.h
include/net-snmp/library/md5.h
include/net-snmp/library/mib.h
include/net-snmp/library/mt_support.h
include/net-snmp/library/oid_stash.h
include/net-snmp/library/parse.h
include/net-snmp/library/read_config.h
include/net-snmp/library/scapi.h
include/net-snmp/library/snmp-tc.h
include/net-snmp/library/snmp.h
include/net-snmp/library/snmpCallbackDomain.h
include/net-snmp/library/snmpTCPDomain.h
include/net-snmp/library/snmpTCPIPv6Domain.h
include/net-snmp/library/snmpUDPDomain.h
include/net-snmp/library/snmpUDPIPv6Domain.h
include/net-snmp/library/snmpUnixDomain.h
include/net-snmp/library/snmp_alarm.h
include/net-snmp/library/snmp_api.h
include/net-snmp/library/snmp_assert.h
include/net-snmp/library/snmp_client.h
include/net-snmp/library/snmp_debug.h
include/net-snmp/library/snmp_enum.h
include/net-snmp/library/snmp_impl.h
include/net-snmp/library/snmp_locking.h
include/net-snmp/library/snmp_logging.h
include/net-snmp/library/snmp_parse_args.h
include/net-snmp/library/snmp_secmod.h
include/net-snmp/library/snmp_transport.h
include/net-snmp/library/snmpusm.h
include/net-snmp/library/snmpv3-security-includes.h
include/net-snmp/library/snmpv3.h
include/net-snmp/library/system.h
include/net-snmp/library/tools.h
include/net-snmp/library/transform_oids.h
include/net-snmp/library/ucd_compat.h
include/net-snmp/library/vacm.h
include/net-snmp/library/winservice.h
include/net-snmp/machine/generic.h
include/net-snmp/mib_api.h
include/net-snmp/net-snmp-config.h
include/net-snmp/net-snmp-includes.h
include/net-snmp/output_api.h
include/net-snmp/pdu_api.h
include/net-snmp/session_api.h
include/net-snmp/snmpv3_api.h
include/net-snmp/system/aix.h
include/net-snmp/system/bsd.h
include/net-snmp/system/bsdi.h
include/net-snmp/system/bsdi3.h
include/net-snmp/system/bsdi4.h
include/net-snmp/system/cygwin.h
include/net-snmp/system/darwin.h
include/net-snmp/system/dynix.h
include/net-snmp/system/freebsd.h
include/net-snmp/system/freebsd2.h
include/net-snmp/system/freebsd3.h
include/net-snmp/system/freebsd4.h
include/net-snmp/system/generic.h
include/net-snmp/system/hpux.h
include/net-snmp/system/irix.h
include/net-snmp/system/linux.h
include/net-snmp/system/mips.h
include/net-snmp/system/netbsd.h
include/net-snmp/system/openbsd.h
include/net-snmp/system/solaris.h
include/net-snmp/system/solaris2.6.h
include/net-snmp/system/solaris2.7.h
include/net-snmp/system/solaris2.8.h
include/net-snmp/system/sunos.h
include/net-snmp/system/svr5.h
include/net-snmp/system/sysv.h
include/net-snmp/system/ultrix4.h
include/net-snmp/types.h
include/net-snmp/utilities.h
include/net-snmp/varbind_api.h
include/net-snmp/version.h
@dirrm include/net-snmp/agent
@dirrm include/net-snmp/library
@dirrm include/net-snmp/machine
@dirrm include/net-snmp/system
@dirrm include/net-snmp
lib/libnetsnmp.a
lib/libnetsnmp.so
lib/libnetsnmp.so.%%shlib%%
lib/libnetsnmpagent.a
lib/libnetsnmpagent.so
lib/libnetsnmpagent.so.%%shlib%%
lib/libnetsnmphelpers.a
lib/libnetsnmphelpers.so
lib/libnetsnmphelpers.so.%%shlib%%
lib/libnetsnmpmibs.a
lib/libnetsnmpmibs.so
lib/libnetsnmpmibs.so.%%shlib%%
sbin/snmpd
sbin/snmptrapd
share/snmp/mib2c.array-user.conf
share/snmp/mib2c.column_defines.conf
share/snmp/mib2c.column_enums.conf
share/snmp/mib2c.conf
share/snmp/mib2c.create-dataset.conf
share/snmp/mib2c.iterate.conf
share/snmp/mib2c.old-api.conf
share/snmp/mibs/AGENTX-MIB.txt
share/snmp/mibs/DISMAN-SCHEDULE-MIB.txt
share/snmp/mibs/DISMAN-SCRIPT-MIB.txt
share/snmp/mibs/EtherLike-MIB.txt
share/snmp/mibs/HCNUM-TC.txt
share/snmp/mibs/HOST-RESOURCES-MIB.txt
share/snmp/mibs/HOST-RESOURCES-TYPES.txt
share/snmp/mibs/IANA-ADDRESS-FAMILY-NUMBERS-MIB.txt
share/snmp/mibs/IANA-LANGUAGE-MIB.txt
share/snmp/mibs/IANAifType-MIB.txt
share/snmp/mibs/IF-INVERTED-STACK-MIB.txt
share/snmp/mibs/IF-MIB.txt
share/snmp/mibs/INET-ADDRESS-MIB.txt
share/snmp/mibs/IP-FORWARD-MIB.txt
share/snmp/mibs/IP-MIB.txt
share/snmp/mibs/IPV6-ICMP-MIB.txt
share/snmp/mibs/IPV6-MIB.txt
share/snmp/mibs/IPV6-TC.txt
share/snmp/mibs/IPV6-TCP-MIB.txt
share/snmp/mibs/IPV6-UDP-MIB.txt
share/snmp/mibs/NET-SNMP-AGENT-MIB.txt
share/snmp/mibs/NET-SNMP-EXAMPLES-MIB.txt
share/snmp/mibs/NET-SNMP-MIB.txt
share/snmp/mibs/NOTIFICATION-LOG-MIB.txt
share/snmp/mibs/RFC-1215.txt
share/snmp/mibs/RFC1155-SMI.txt
share/snmp/mibs/RFC1213-MIB.txt
share/snmp/mibs/RMON-MIB.txt
share/snmp/mibs/SMUX-MIB.txt
share/snmp/mibs/SNMP-COMMUNITY-MIB.txt
share/snmp/mibs/SNMP-FRAMEWORK-MIB.txt
share/snmp/mibs/SNMP-MPD-MIB.txt
share/snmp/mibs/SNMP-NOTIFICATION-MIB.txt
share/snmp/mibs/SNMP-PROXY-MIB.txt
share/snmp/mibs/SNMP-TARGET-MIB.txt
share/snmp/mibs/SNMP-USER-BASED-SM-MIB.txt
share/snmp/mibs/SNMP-VIEW-BASED-ACM-MIB.txt
share/snmp/mibs/SNMPv2-CONF.txt
share/snmp/mibs/SNMPv2-MIB.txt
share/snmp/mibs/SNMPv2-SMI.txt
share/snmp/mibs/SNMPv2-TC.txt
share/snmp/mibs/SNMPv2-TM.txt
share/snmp/mibs/TCP-MIB.txt
share/snmp/mibs/UCD-DEMO-MIB.txt
share/snmp/mibs/UCD-DISKIO-MIB.txt
share/snmp/mibs/UCD-DLMOD-MIB.txt
share/snmp/mibs/UCD-IPFWACC-MIB.txt
share/snmp/mibs/UCD-SNMP-MIB.txt
share/snmp/mibs/UDP-MIB.txt
@unexec rm -f %D/share/snmp/mibs/.index 2>/dev/null || true
share/snmp/snmpconf-data/snmp-data/authopts
share/snmp/snmpconf-data/snmp-data/debugging
share/snmp/snmpconf-data/snmp-data/mibs
share/snmp/snmpconf-data/snmp-data/output
share/snmp/snmpconf-data/snmp-data/snmpconf-config
share/snmp/snmpconf-data/snmpd-data/acl
share/snmp/snmpconf-data/snmpd-data/basic_setup
share/snmp/snmpconf-data/snmpd-data/extending
share/snmp/snmpconf-data/snmpd-data/monitor
share/snmp/snmpconf-data/snmpd-data/operation
share/snmp/snmpconf-data/snmpd-data/snmpconf-config
share/snmp/snmpconf-data/snmpd-data/system
share/snmp/snmpconf-data/snmpd-data/trapsinks
share/snmp/snmpconf-data/snmptrapd-data/formatting
share/snmp/snmpconf-data/snmptrapd-data/snmpconf-config
share/snmp/snmpconf-data/snmptrapd-data/traphandle
@dirrm share/snmp/snmpconf-data/snmp-data
@dirrm share/snmp/snmpconf-data/snmpd-data
@dirrm share/snmp/snmpconf-data/snmptrapd-data
@dirrm share/snmp/snmpconf-data
@dirrm share/snmp/mibs
@unexec rmdir %D/share/snmp 2>/dev/null || true
