# security/bitwarden_rs/files/bitwarden_rs.in
#!/bin/sh

#
# $FreeBSD$
#

# PROVIDE: bitwarden_rs
# REQUIRE: LOGIN DAEMON NETWORKING
# KEYWORD: rust

# Enable this script by adding:
# bitwarden_rs_enable="YES"
# ... to /etc/rc.conf
# rc_debug=yes

. /etc/rc.subr
name=bitwarden_rs

rcvar=bitwarden_rs_enable
load_rc_config ${name}

: ${bitwarden_rs_chdir=%%PREFIX%%/www/bitwarden_rs}
: ${bitwarden_rs_user:=www}
: ${bitwarden_rs_group:=www}

# This is the tool init launches
command="/usr/sbin/daemon"

pidfile="/var/run/${name}/${name}.pid"

task="%%PREFIX%%/bin/${name}"
procname="${task}"

command_args="-p ${pidfile} -T ${name} ${task}"

start_precmd="start_precmd"
start_precmd()
{
        if [ ! -e "/var/run/${name}" ] ; then
                install -d -o ${bitwarden_rs_user} -g ${bitwarden_rs_group} /var/run/${name};
        fi

        # checkyesno rubywarden_api_signups
        # if [ "$?" -eq 0 ]; then
        #         export ALLOW_SIGNUPS=1
        #         echo "Bitwarden Signups Enabled"
        # fi
}

run_rc_command "$1"
