# sysutils/moosefs3-master/files/openrc-mfscgiserv.in
#!/sbin/openrc-run

name=mfscgiserv

command=%%PREFIX%%/sbin/${name}
pidfile=%%MFS_WORKDIR%%/mfs/.${name}.lock
procname="python"

stop_pre() {
	$command stop
}

depend() {
	keyword -shutdown
}
