# databases/pgbouncer/files/pgbouncer.sh.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: pgbouncer
# REQUIRE: LOGIN
# BEFORE:  securelevel
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable pgbouncer:
#
# pgbouncer_enable (bool):      Set to "YES" to enable pgbouncer
#                               Default is "NO".
# pgbouncer_conf (path):        Set full path to configuration file.
#                               Default is "%%PREFIX%%/etc/pgbouncer.ini".
# pgbouncer_flags (string):     Additional flags passed to pgbouncer.
#                               Default is "".

. /etc/rc.subr

name="pgbouncer"
rcvar=pgbouncer_enable

load_rc_config "$name"
: ${pgbouncer_enable="NO"}
: ${pgbouncer_user="pgbouncer"}
: ${pgbouncer_conf="%%PREFIX%%/etc/$name.ini"}

extra_commands="reload"

command="%%PREFIX%%/bin/pgbouncer"
pidfile="%%PGBOUNCER_RUNDIR%%/$name.pid"
required_files="${pgbouncer_conf}"
command_args="-d ${pgbouncer_conf}"

run_rc_command "$1"
