# net/openvswitch/files/ovs-vswitchd.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ovs_vswitchd
# REQUIRE: LOGIN ovsdb_server
# KEYWORD: shutdown

. /etc/rc.subr

name="ovs_vswitchd"
rcvar=ovs_vswitchd_enable

command="/usr/local/sbin/ovs-vswitchd"
start_precmd="ovs_vswitchd_precmd"

load_rc_config $name

: ${ovs_vswitchd_enable="NO"}

ovs_vswitchd_precmd() {
	load_kld if_tap || return 1
}

command_args="--detach --pidfile"

run_rc_command "$1"
