# net/vtun/files/openrc-vtunclient.in
#!/sbin/openrc-run
#
# vtunclient_flags="<session> <server address>"
#
# If you want to use multiple connections, use $vtunclient_list.
# Usual $vtunclient_flags will be ignored when $vtunclient_list is
# set.
#
# vtunclient_list="site1 site2"
# vtunclient_site1_flags="<session> <server address>"
# vtunclient_site1_pidfile="/var/run/vtun-site1.pid"
# vtunclient_site2_flags="<session> <server address>"
# vtunclient_site2_pidfile="/var/run/vtun-site2.pid"
# ...
#

name=vtunclient
required_files="%%PREFIX%%/etc/vtund.conf"

command=%%PREFIX%%/sbin/vtund

start_pre() {
	if [ -n "$*" ]; then
		vtunclient_list="$*"
	fi

	if [ -n "$vtunclient_list" ]; then
		for _client in ${vtunclient_list}; do
			eval command_args=\$vtunclient_${_client}_flags
			eval pidfile=\$vtunclient_${_client}_pidfile
			command_args="-z $pidfile $command_args"
			start-stop-daemon --start --name ${name} --pidfile ${pidfile} --exec ${command} ${vtunclient_flags}
		done
	fi
}

stop_pre() {
	if [ -n "$*" ]; then
		vtunclient_list="$*"
	fi

	if [ -n "$vtunclient_list" ]; then
		for _client in ${vtunclient_list}; do
			eval command_args=\$vtunclient_${_client}_flags
			eval pidfile=\$vtunclient_${_client}_pidfile
			command_args="-z $pidfile $command_args"
			start-stop-daemon -- ${pidfile} -- stop $command
		done
	fi
}
