# net/vtun/files/vtunclient.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: vtunclient
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable vtunclient:
#
# vtunclient_enable="YES"
# vtunclient_flags="<session> <server address>"
#
# If you want to use multiple connections, use $vtunclient_list.
# Usual $vtunclient_flags will be ignored when $vtunclient_list is
# set.
#
# vtunclient_enable="YES"
# vtunclient_list="site1 site2"
# vtunclient_site1_flags="<session> <server address>"
# vtunclient_site1_pidfile="/var/run/vtun-site1.pid"
# vtunclient_site2_flags="<session> <server address>"
# vtunclient_site2_pidfile="/var/run/vtun-site2.pid"
# ...
#

vtunclient_enable=${vtunclient_enable:-"NO"}

. /etc/rc.subr

name=vtunclient
rcvar=`set_rcvar`
required_files="%%PREFIX%%/etc/vtund.conf"

command=%%PREFIX%%/sbin/vtund

load_rc_config $name
cmd="$1"
if [ $# -gt 0 ]; then
	shift
fi
if [ -n "$*" ]; then
	vtunclient_list="$*"
fi

if [ -z "$vtunclient_list" ]; then
	run_rc_command "$cmd"
else
	for _client in ${vtunclient_list}; do
		eval vtunclient_flags=\$vtunclient_${_client}_flags
		eval pidfile=\$vtunclient_${_client}_pidfile
		run_rc_command "$cmd"
	done
fi
