# databases/courier-authlib-pgsql/pkg-plist
bin/courierauthconfig
etc/rc.d/courier-authdaemond.sh
%%MYSQLFLAG%%etc/authlib/authmysqlrc.dist
%%PGSQLFLAG%%etc/authlib/authpgsqlrc.dist
%%LDAPFLAG%%etc/authlib/authldaprc.dist
etc/authlib/authdaemonrc.dist
include/courier_auth_config.h
include/courierauth.h
include/courierauthsasl.h
include/courierauthsaslclient.h
include/courierauthdebug.h
lib/courier-authlib/libcourierauth.so.0
lib/courier-authlib/libcourierauth.so
lib/courier-authlib/libcourierauth.a
lib/courier-authlib/libcourierauthsasl.so.0
lib/courier-authlib/libcourierauthsasl.so
lib/courier-authlib/libcourierauthsasl.a
lib/courier-authlib/libcourierauthsaslclient.so.0
lib/courier-authlib/libcourierauthsaslclient.so
lib/courier-authlib/libcourierauthsaslclient.a
lib/courier-authlib/libcourierauthcommon.so.0
lib/courier-authlib/libcourierauthcommon.so
lib/courier-authlib/libcourierauthcommon.a
lib/courier-authlib/libauthcustom.so.0
lib/courier-authlib/libauthcustom.so
lib/courier-authlib/libauthuserdb.so.0
lib/courier-authlib/libauthuserdb.so
%%PAMFLAG%%lib/courier-authlib/libauthpam.so.0
%%PAMFLAG%%lib/courier-authlib/libauthpam.so
%%MYSQLFLAG%%lib/courier-authlib/libauthmysql.so.0
%%MYSQLFLAG%%lib/courier-authlib/libauthmysql.so
%%PGSQLFLAG%%lib/courier-authlib/libauthpgsql.so.0
%%PGSQLFLAG%%lib/courier-authlib/libauthpgsql.so
%%LDAPFLAG%%lib/courier-authlib/libauthldap.so.0
%%LDAPFLAG%%lib/courier-authlib/libauthldap.so
%%VPOPMAILFLAG%%lib/courier-authlib/libauthvchkpw.so.0
%%VPOPMAILFLAG%%lib/courier-authlib/libauthvchkpw.so
libexec/courier-authlib/authsystem.passwd
libexec/courier-authlib/authdaemond
libexec/courier-authlib/makedatprog
sbin/authtest
sbin/authdaemond
sbin/authenumerate
sbin/courierlogger
sbin/makeuserdb
sbin/userdb
sbin/userdb-test-cram-md5
sbin/userdbpw
sbin/vchkpw2userdb
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@unexec rmdir %D/var/authdaemon 2>/dev/null || true
@unexec rmdir %D/var 2>/dev/null || true
@dirrm etc/authlib
@unexec rmdir %D/libexec/authlib 2>/dev/null || true
