# devel/pcsc-lite/files/openrc-pcscd.in
#!/sbin/openrc-run

name="pcscd"
command="%%PREFIX%%/sbin/pcscd"
