# x11/kde4-workspace/files/kdm4.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: kdm4
# REQUIRE: LOGIN cleanvar moused syscons dbus hald
#
# Add the following to /etc/rc.conf to start KDM at boot time:
#
# kdm4_enable="YES"
#

. /etc/rc.subr

kdm4_enable=${kdm4_enable-"NO"}

GENKDMCONF=%%PREFIX%%/bin/genkdmconf
KDMCONFDIR=%%PREFIX%%/share/config/kdm

name=kdm4
rcvar=`set_rcvar`
command="%%PREFIX%%/bin/kdm-bin"
pidfile="/var/run/kdm.pid"
start_cmd="kdm_start"

kdm_start()
{
    if ! checkyesno kdm4_enable; then
        return 0
    fi
    echo "Starting ${name}."

    # Configure KDM if needed.
    if [ ! -r ${KDMCONFDIR}/kdmrc ]; then
        echo "Generating KDM configuration."
        ${GENKDMCONF} --no-old --in ${KDMCONFDIR}
    else
        echo "Updating KDM configuration."
        ${GENKDMCONF} --in ${KDMCONFDIR}
    fi

    ( iter=0
    while ! pgrep -f "^/usr/libexec/getty " > /dev/null 2>&1; do
        if [ ${iter} -ge 600 ]; then
            return 1
        fi
        sleep 1
        iter=$(expr ${iter} + 1)
    done
    if checkyesno hald_enable; then
        iter=0
        while [ ${iter} -lt 60 ] &&
          ! %%LOCALBASE%%/bin/lshal > /dev/null 2>&1; do
            sleep 1
            iter=$(expr ${iter} + 1)
        done
    fi
    ${command} ${kdm4_flags} ) &
}

load_rc_config ${name}
run_rc_command "$1"
