# japanese/Canna/pkg-plist
bin/addwords
bin/canlisp
bin/cannacheck
bin/cannastat
bin/catdic
bin/chkconc
bin/chmoddic
bin/cpdic
bin/crfreq
bin/crxdic
bin/crxgram
bin/cshost
bin/ctow
bin/delwords
bin/dicar
bin/dpbindic
bin/dpromdic
bin/dpxdic
bin/forcpp
bin/forsort
bin/kpdic
bin/lsdic
bin/mergeword
bin/mkbindic
bin/mkdic
bin/mkromdic
bin/mvdic
bin/rmdic
bin/splitword
bin/syncdic
bin/wtoc
etc/rc.d/canna.sh
include/canna/RK.h
include/canna/jrkanji.h
include/canna/keydef.h
include/canna/mfdef.h
include/canna/pubconf.h
include/canna/sysdep.h
share/canna/default.canna
share/canna/dic/canna/bushu.cld
share/canna/dic/canna/chimei.ctd
share/canna/dic/canna/dics.dir
share/canna/dic/canna/fuzokugo.cbd
share/canna/dic/canna/hojomwd.ctd
share/canna/dic/canna/hojoswd.ctd
share/canna/dic/canna/iroha.cbd
share/canna/dic/canna/iroha.cld
share/canna/dic/canna/kanasmpl.ctd
share/canna/dic/canna/katakana.ctd
share/canna/dic/canna/keishiki.ctd
share/canna/dic/canna/necgaiji.ctd
share/canna/dic/canna/number.ctd
share/canna/dic/canna/software.ctd
share/canna/dic/canna/suffix.ctd
share/canna/dic/default.cbp
share/canna/dic/english.cbp
share/canna/dic/jdaemon.cbp
share/canna/dic/just.cbp
share/canna/dic/kaisoku.cbp
share/canna/dic/kana.cbp
share/canna/dic/lan5.cbp
share/canna/dic/matsu.cbp
share/canna/dic/newjis.cbp
share/canna/dic/romaji.cbp
share/canna/dic/sokuon.cbp
share/canna/dic/tut.cbp
share/canna/dic/vje.cbp
share/canna/dic/wx2+.cbp
share/canna/sample/1.1.canna
share/canna/sample/1.2.canna
share/canna/sample/fixer.canna
share/canna/sample/jdemon.canna
share/canna/sample/just.canna
share/canna/sample/lan5.canna
share/canna/sample/matsu.canna
share/canna/sample/modes.canna
share/canna/sample/sample.canna
share/canna/sample/skk.canna
share/canna/sample/src/default.ctd
share/canna/sample/src/english.ctd
share/canna/sample/src/is.c
share/canna/sample/src/jdaemon.ctd
share/canna/sample/src/just.ctd
share/canna/sample/src/kaisoku.ctd
share/canna/sample/src/kana.ctd
share/canna/sample/src/lan5.ctd
share/canna/sample/src/matsu.ctd
share/canna/sample/src/newjis.ctd
share/canna/sample/src/romaji.ctd
share/canna/sample/src/sokuon.ctd
share/canna/sample/src/tut.ctd
share/canna/sample/src/vje.ctd
share/canna/sample/src/wx2+.ctd
share/canna/sample/tut.canna
share/canna/sample/unix.canna
share/canna/sample/verbose.canna
share/canna/sample/vje.canna
share/canna/sample/wx2+.canna
lib/libRKC.a
lib/libRKC.so
lib/libRKC.so.1
lib/libRKC16.a
lib/libRKC16.so
lib/libRKC16.so.1
lib/libcanna.a
lib/libcanna.so
lib/libcanna.so.1
lib/libcanna16.a
lib/libcanna16.so
lib/libcanna16.so.1
sbin/cannaserver
sbin/cannakill
@exec mkdir -p %D/share/canna/dic/user/user %D/share/canna/dic/group
@exec chown bin:bin %D/share/canna/dic/user/user %D/share/canna/dic/user %D/share/canna/dic/group
@exec chmod 775 %D/share/canna/dic/user/user %D/share/canna/dic/user %D/share/canna/dic/group
@dirrm share/canna/sample/src
@dirrm share/canna/sample
@dirrm share/canna/dic/canna
@dirrm include/canna
@dirrm lib/canna/dic/user/user
@dirrm lib/canna/dic/user
@dirrm lib/canna/dic/group
@dirrm lib/canna/dic
@dirrm lib/canna
@unexec rmdir %D/share/canna/dic/user/user 2>/dev/null || true
@unexec rmdir %D/share/canna/dic/user      2>/dev/null || true
@unexec rmdir %D/share/canna/dic/group     2>/dev/null || true
@unexec rmdir %D/share/canna/dic           2>/dev/null || true
@unexec rmdir %D/share/canna               2>/dev/null || true
