# japanese/Canna/pkg-plist
bin/addwords
bin/canlisp
bin/cannacheck
bin/cannastat
bin/canvert
bin/catdic
bin/chmoddic
bin/cpdic
bin/crfreq
bin/crrdic
bin/crxdic
bin/crxgram
bin/cshost
bin/ctow
bin/delwords
bin/dicar
bin/dpbindic
bin/dpromdic
bin/dpwdic
bin/dpxdic
bin/forcpp
bin/forsort
bin/itoc
bin/kpdic
bin/lsdic
bin/mergeword
bin/mkbindic
bin/mkdic
bin/mkromdic
bin/mvdic
bin/rmdic
bin/splitword
bin/syncdic
bin/wtoc
etc/rc.d/canna.sh
include/canna/RK.h
include/canna/jrkanji.h
include/canna/keydef.h
include/canna/mfdef.h
lib/canna/default.canna
lib/canna/dic/canna/bushu.fq
lib/canna/dic/canna/dics.dir
lib/canna/dic/canna/fuzokugo.d
lib/canna/dic/canna/fuzokugo.fq
lib/canna/dic/canna/hojomwd.t
lib/canna/dic/canna/hojoswd.t
lib/canna/dic/canna/iroha.d
lib/canna/dic/canna/iroha.fq
lib/canna/dic/canna/number.t
lib/canna/dic/canna/suffix.t
lib/canna/dic/canna/yuubin.fq
lib/canna/dic/default.kp
lib/canna/dic/default.rdic
lib/canna/dic/english.kp
lib/canna/dic/hyoujun.gram
lib/canna/dic/jdaemon.kp
lib/canna/dic/jdaemon.rdic
lib/canna/dic/just.kp
lib/canna/dic/just.rdic
lib/canna/dic/kaisoku.kp
lib/canna/dic/kana.kp
lib/canna/dic/kana.rdic
lib/canna/dic/lan5.kp
lib/canna/dic/newjis.kp
lib/canna/dic/newjis.rdic
lib/canna/dic/premountdics
lib/canna/dic/romaji.kp
lib/canna/dic/sokuon.kp
lib/canna/dic/tut.kp
lib/canna/dic/user/user/chimei.t
lib/canna/dic/user/user/dics.dir
lib/canna/dic/user/user/katakana.t
lib/canna/dic/user/user/necgaiji.t
lib/canna/dic/user/user/software.t
lib/canna/dic/user/user/user.t
lib/canna/engine.cf
lib/canna/sample/1.1.canna
lib/canna/sample/1.1.iroha
lib/canna/sample/1.2.canna
lib/canna/sample/default.canna
lib/canna/sample/default.iroha
lib/canna/sample/fixer.canna
lib/canna/sample/fixer.iroha
lib/canna/sample/jdemon.canna
lib/canna/sample/jdemon.iroha
lib/canna/sample/just.canna
lib/canna/sample/just.iroha
lib/canna/sample/lan5.canna
lib/canna/sample/modes.canna
lib/canna/sample/modes.iroha
lib/canna/sample/sample.canna
lib/canna/sample/sample.iroha
lib/canna/sample/src/Makefile
lib/canna/sample/src/default.kpdef
lib/canna/sample/src/english.kpdef
lib/canna/sample/src/is.c
lib/canna/sample/src/jdaemon.kpdef
lib/canna/sample/src/just.kpdef
lib/canna/sample/src/kaisoku.kpdef
lib/canna/sample/src/kana.kpdef
lib/canna/sample/src/lan5.kpdef
lib/canna/sample/src/newjis.kpdef
lib/canna/sample/src/romaji.kpdef
lib/canna/sample/src/sokuon.kpdef
lib/canna/sample/src/tut.kpdef
lib/canna/sample/tut.canna
lib/canna/sample/unix.canna
lib/canna/sample/vje.canna
lib/canna/sample/vje.iroha
lib/libRKC.a
lib/libRKC.so.1.0
lib/libRKC16.a
lib/libRKC16.so.1.0
lib/libcanna.a
lib/libcanna.so.1.0
lib/libcanna16.a
lib/libcanna16.so.1.0
@exec ldconfig -m %B
@unexec /sbin/ldconfig -R
man/man1/addwords.1.gz
man/man1/cannacheck.1.gz
man/man1/cannakill.1.gz
man/man1/cannaserver.1.gz
man/man1/cannastat.1.gz
man/man1/canvert.1.gz
man/man1/catdic.1.gz
man/man1/cpdic.1.gz
man/man1/cshost.1.gz
man/man1/ctow.1.gz
man/man1/delwords.1.gz
man/man1/dicar.1.gz
man/man1/dpbindic.1.gz
man/man1/dpromdic.1.gz
man/man1/lsdic.1.gz
man/man1/mkbindic.1.gz
man/man1/mkdic.1.gz
man/man1/mkromdic.1.gz
man/man1/mvdic.1.gz
man/man1/rmdic.1.gz
man/man1/splitword.1.gz
man/man1/wtoc.1.gz
man/man3/RkCloseRoma.3.gz
man/man3/RkCvtEuc.3.gz
man/man3/RkCvtHan.3.gz
man/man3/RkCvtHira.3.gz
man/man3/RkCvtKana.3.gz
man/man3/RkCvtRoma.3.gz
man/man3/RkCvtZen.3.gz
man/man3/RkDefineDic.3.gz
man/man3/RkDeleteDic.3.gz
man/man3/RkEndBun.3.gz
man/man3/RkEnlarge.3.gz
man/man3/RkFinalize.3.gz
man/man3/RkGetDicList.3.gz
man/man3/RkGetKanji.3.gz
man/man3/RkGetLex.3.gz
man/man3/RkGetStat.3.gz
man/man3/RkGetYomi.3.gz
man/man3/RkGoto.3.gz
man/man3/RkInitialize.3.gz
man/man3/RkIntro.3.gz
man/man3/RkLeft.3.gz
man/man3/RkMapPhonogram.3.gz
man/man3/RkMapRoma.3.gz
man/man3/RkMountD.3.gz
man/man3/RkNext.3.gz
man/man3/RkNfer.3.gz
man/man3/RkOpenRoma.3.gz
man/man3/RkPrev.3.gz
man/man3/RkRemountDic.3.gz
man/man3/RkResize.3.gz
man/man3/RkRgnBun.3.gz
man/man3/RkRight.3.gz
man/man3/RkShorten.3.gz
man/man3/RkStoreYomi.3.gz
man/man3/RkUnmountDic.3.gz
man/man3/RkXfer.3.gz
man/man3/cannalib.3.gz
man/man3/uilib.3.gz
sbin/cannakill
sbin/cannaserver
@exec mkdir /var/run/canna
@exec chown bin.bin /var/run/canna
@exec chown bin.bin %D/lib/canna/dic/canna
@exec chmod 775 %D/lib/canna/dic/canna
@exec chown bin.bin %D/lib/canna/dic
@exec chmod 775 %D/lib/canna/dic
@exec chown bin.bin %D/lib/canna/dic/user/user %D/lib/canna/dic/user
@exec chmod 775 %D/lib/canna/dic/user/user %D/lib/canna/dic/user
@unexec /bin/rmdir /var/run/canna
