# devel/p5-Gearman-Server/files/openrc-gearmand-perl.in
#!/sbin/openrc-run

name="%%NAME%%"
command="%%PREFIX%%/bin/gearmand"
command_args="--daemon"
