# net/isc-dhcp3-server/pkg-plist
@comment $FreeBSD$
etc/dhclient-script
sbin/dhclient
sbin/dhcpd
sbin/dhcrelay
etc/rc.d/isc-dhcpd.sh.sample
include/isc/boolean.h
include/isc/dst.h
include/isc/int.h
include/isc/lang.h
include/isc/list.h
include/isc/result.h
include/isc/types.h
@dirrm include/isc
include/omapip/alloc.h
include/omapip/buffer.h
include/omapip/omapip.h
@dirrm include/omapip
include/dhcpctl.h
lib/libomapi.a
lib/libdhcpctl.a
@exec touch /var/db/dhcpd.leases
@exec touch %D/etc/dhclient.conf
