# net-mgmt/p5-FusionInventory-Agent/files/openrc-fusioninventory.in
#!/sbin/openrc-run

name=fusioninventory
pidfile=/var/run/${name}.pid

command=%%PREFIX%%/bin/fusioninventory-agent
command_args="--daemon --pidfile=${pidfile}"

depend() {
	keyword -shutdown
}
