# databases/mysql54-server/pkg-plist
bin/comp_err
bin/isamchk
bin/isamlog
bin/msql2mysql
bin/mysql
bin/mysql_fix_privilege_tables
bin/mysql_install_db
bin/mysql_setpermission
bin/mysql_zap
bin/mysqlaccess
bin/mysqladmin
bin/mysqlbug
bin/mysqldump
bin/mysqlimport
bin/mysqlshow
bin/perror
bin/replace
bin/resolveip
bin/safe_mysqld
etc/rc.d/mysql.sh
include/mysql/dbug.h
include/mysql/errmsg.h
include/mysql/m_ctype.h.in
include/mysql/m_string.h
include/mysql/my_sys.h
include/mysql/mysql.h
include/mysql/mysql_com.h
include/mysql/mysql_version.h
include/mysql/mysqld_error.h
@unexec install-info --delete %D/info/mysql.info %D/info/dir
info/mysql.info
@exec [ -f %D/info/dir -o ! -f /usr/share/info/dir ] || sed -ne '1,/Menu:/p' /usr/share/info/dir > %D/info/dir
@exec install-info %D/info/mysql.info %D/info/dir
lib/mysql/libdbug.a
lib/mysql/libheap.a
lib/mysql/libmerge.a
lib/mysql/libmysqlclient.a
lib/mysql/libmysqlclient.la
lib/mysql/libmysqlclient.so
lib/mysql/libmysqlclient.so.6
lib/mysql/libmystrings.a
lib/mysql/libmysys.a
lib/mysql/libnisam.a
libexec/mysqld
share/doc/mysql/include.texi
share/doc/mysql/manual.html
share/doc/mysql/manual.ps
share/doc/mysql/manual.texi
share/doc/mysql/manual.txt
share/doc/mysql/manual_toc.html
share/mysql/binary-configure
share/mysql/czech/errmsg.sys
share/mysql/czech/errmsg.txt
share/mysql/dutch/errmsg.sys
share/mysql/dutch/errmsg.txt
share/mysql/english/errmsg.sys
share/mysql/english/errmsg.txt
share/mysql/french/errmsg.sys
share/mysql/french/errmsg.txt
share/mysql/german/errmsg.sys
share/mysql/german/errmsg.txt
share/mysql/hungarian/errmsg.sys
share/mysql/hungarian/errmsg.txt
share/mysql/italian/errmsg.sys
share/mysql/italian/errmsg.txt
share/mysql/japanese/errmsg.sys
share/mysql/japanese/errmsg.txt
share/mysql/make_binary_distribution
share/mysql/my-example.cnf
share/mysql/mysql-3.22.16a-gamma.spec
share/mysql/mysql-log-rotate
share/mysql/mysql.server
share/mysql/norwegian-ny/errmsg.sys
share/mysql/norwegian-ny/errmsg.txt
share/mysql/norwegian/errmsg.sys
share/mysql/norwegian/errmsg.txt
share/mysql/polish/errmsg.sys
share/mysql/polish/errmsg.txt
share/mysql/portuguese/errmsg.sys
share/mysql/portuguese/errmsg.txt
share/mysql/russian/errmsg.sys
share/mysql/russian/errmsg.txt
share/mysql/slovak/errmsg.sys
share/mysql/slovak/errmsg.txt
share/mysql/spanish/errmsg.sys
share/mysql/spanish/errmsg.txt
share/mysql/swedish/errmsg.sys
share/mysql/swedish/errmsg.txt
@dirrm include/mysql
@dirrm lib/mysql
@dirrm share/doc/mysql
@dirrm share/mysql/czech
@dirrm share/mysql/dutch
@dirrm share/mysql/english
@dirrm share/mysql/french
@dirrm share/mysql/german
@dirrm share/mysql/hungarian
@dirrm share/mysql/italian
@dirrm share/mysql/japanese
@dirrm share/mysql/mysql
@dirrm share/mysql/norwegian
@dirrm share/mysql/norwegian-ny
@dirrm share/mysql/polish
@dirrm share/mysql/portuguese
@dirrm share/mysql/russian
@dirrm share/mysql/slovak
@dirrm share/mysql/spanish
@dirrm share/mysql/swedish
@dirrm share/mysql
@exec %D/bin/mysql_install_db
@exec /usr/bin/killall mysqld
@exec /usr/bin/env OBJFORMAT=%%PORTOBJFORMAT%% /sbin/ldconfig -m %D/lib/mysql
@unexec /usr/bin/env OBJFORMAT=%%PORTOBJFORMAT%% /sbin/ldconfig -R
