# misc/yaunc/pkg-plist
bin/yaunc
etc/rc.d/yaunc.sh
etc/yaunc.conf.sample
