# net/msend/files/openrc-mesgd.in
#!/sbin/openrc-run

name="mesgd"
command=%%PREFIX%%/sbin/${name}

depend() {
	keyword -shutdown
}
