# net-mgmt/coovachilli/files/openrc-coovachilli.in
#!/sbin/openrc-run

name=coovachilli

command=%%PREFIX%%/sbin/chilli

depend() {
	uset netif
	keyword -jail
}
