# mail/exilog/pkg-plist
%%AGENT%%@unexec [ ! -f /var/run/exilog.pid ] || %D/etc/rc.d/%%RCSCRIPT%% stop
%%AGENT%%sbin/exilog_agent.pl
sbin/exilog_cleanup.pl
%%EXILOGDIR%%/exilog_cgi.pl
%%EXILOGDIR%%/exilog_cgi_html.pm
%%EXILOGDIR%%/exilog_cgi_messages.pm
%%EXILOGDIR%%/exilog_cgi_param.pm
%%EXILOGDIR%%/exilog_cgi_queues.pm
%%EXILOGDIR%%/exilog_cgi_servers.pm
%%EXILOGDIR%%/exilog_config.pm
%%EXILOGDIR%%/exilog_parse.pm
%%EXILOGDIR%%/exilog_sql.pm
%%EXILOGDIR%%/exilog_util.pm
%%EXILOGDIR%%/exilog_jscript.js
%%EXILOGDIR%%/exilog_stylesheet.css
%%EXILOGDIR%%/.htaccess
%%EXILOGDIR%%/icons/address.png
%%EXILOGDIR%%/icons/arrival.png
%%EXILOGDIR%%/icons/arrival_auth.png
%%EXILOGDIR%%/icons/arrival_local.png
%%EXILOGDIR%%/icons/arrival_normal.png
%%EXILOGDIR%%/icons/arrival_tls.png
%%EXILOGDIR%%/icons/arrival_tls_auth.png
%%EXILOGDIR%%/icons/deferral_normal.png
%%EXILOGDIR%%/icons/deferral_tls.png
%%EXILOGDIR%%/icons/deferred.png
%%EXILOGDIR%%/icons/delivered.png
%%EXILOGDIR%%/icons/delivery.png
%%EXILOGDIR%%/icons/delivery_normal.png
%%EXILOGDIR%%/icons/delivery_tls.png
%%EXILOGDIR%%/icons/dns.png
%%EXILOGDIR%%/icons/dsn_warning.png
%%EXILOGDIR%%/icons/errmsg.png
%%EXILOGDIR%%/icons/error.png
%%EXILOGDIR%%/icons/error_normal.png
%%EXILOGDIR%%/icons/error_tls.png
%%EXILOGDIR%%/icons/event_type.png
%%EXILOGDIR%%/icons/find.png
%%EXILOGDIR%%/icons/frozen.png
%%EXILOGDIR%%/icons/helo.png
%%EXILOGDIR%%/icons/ident.png
%%EXILOGDIR%%/icons/queue_deferred.png
%%EXILOGDIR%%/icons/queue_frozen.png
%%EXILOGDIR%%/icons/queue_normal.png
%%EXILOGDIR%%/icons/queued.png
%%EXILOGDIR%%/icons/reject_postdata.png
%%EXILOGDIR%%/icons/reject_predata.png
%%EXILOGDIR%%/icons/router_transport.png
%%EXILOGDIR%%/icons/server.png
%%EXILOGDIR%%/icons/server_normal.png
%%EXILOGDIR%%/icons/size.png
%%EXILOGDIR%%/icons/stats_h24.png
%%EXILOGDIR%%/icons/stopwatch.png
%%EXILOGDIR%%/icons/timerange.png
%%EXILOGDIR%%/icons/unknown.png
etc/exilog.conf-dist
%%PORTDOCS%%%%DOCSDIR%%/Changelog
%%PORTDOCS%%%%DOCSDIR%%/exilog.txt
%%PORTDOCS%%%%DOCSDIR%%/mysql-db-script.sql
%%PORTDOCS%%%%DOCSDIR%%/pgsql-db-script.sql
@dirrm %%EXILOGDIR%%/icons
@dirrm %%EXILOGDIR%%
%%PORTDOCS%%@dirrm %%DOCSDIR%%
