# misc/e2fsprogs-libuuid/files/openrc-uuidd.in
#!/sbin/openrc-run
#
# rcNG script to start uuidd at boot-time on rcNG-enabled systems,
# such as FreeBSD.  Note: Starting uuidd at boot-time is not strictly 
# necessary, the library will - as of 1.41.5 - silently launch an 
# instance of uuidd that exits after 300 seconds; for most accurate 
# time-based uuids generated from unprivileged user accounts it may be 
# useful to run it system-wide.
#
# (C) 2008, 2009 by Matthias Andree.
# Licensed under the modified (= 2-clause) BSD license.

name="uuidd"
command="%%PREFIX%%/sbin/uuidd"
