# net/freeswitch/files/openrc-freeswitch.in
#!/sbin/openrc-run

name=freeswitch

: ${freeswitch_user:=freeswitch}
: ${freeswitch_group:=freeswitch}

command=%%PREFIX%%/bin/freeswitch
command_args="-nc"

depend() {
	keyword -shutdown
}
