# net-mgmt/bpft/pkg-plist
bin/trafdump
bin/traflog
bin/trafsave
bin/trafstart
bin/trafstat
bin/trafstop
etc/rc.d/trafd.sh
etc/traflog.format
sbin/trafd
sbin/trafstatd
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/WhatsNew
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrmtry var/bpft
@dirrmtry var
