# net-mgmt/bpft/pkg-plist
sbin/trafd
sbin/trafstatd
bin/trafdump
bin/traflog
bin/trafsave
bin/trafstart
bin/trafstat
bin/trafstop
etc/rc.d/trafd.sh
etc/traflog.format
@dirrmtry var/bpft
@dirrmtry var
