# irc/minbif/files/minbif.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: minbif
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable minbif:
# minbif_enable="YES"

. /etc/rc.subr 

name="minbif"
rcvar=`set_rcvar`

load_rc_config $name

: ${minbif_enable="NO"}
: ${minbif_user="%%USERS%%"}
: ${minbif_config="%%PREFIX%%/etc/minbif/minbif.conf"}

required_files="${minbif_config}"
command="%%PREFIX%%/bin/minbif"
command_args="${minbif_config}"

minbif_start() {
	su -fm ${minbif_user} -c "${command} start"
}

minbif_stop() {
	su -fm ${minbif_user} -c "${command} stop"
}

run_rc_command "$1"
