# security/f-protd/files/fprotd_scanmail.sh.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: fprotd_scanmail
# REQUIRE: fprotd
# KEYWORD: shutdown
#
# Add the following to /etc/rc.conf[.local] to enable this service
#
# fprotd_scanmail_enable="YES"
#

. /etc/rc.subr

name=fprotd_scanmail
rcvar=`set_rcvar`
pidfile="/var/run/scanmail.pid"
command="%%PREFIX%%/f-prot/tools/scan-mail.pl"
command_args="> /dev/null 2>&1"
command_interpreter="/usr/bin/perl"

required_files="/etc/f-prot.conf"

load_rc_config ${name}

: ${fprotd_scanmail_enable=NO}
: ${fprotd_scanmail_flags="--daemon --milter inet:12200@127.0.0.1"}

run_rc_command "$1"
