# www/varnish2/files/varnishncsa.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: varnishncsa
# REQUIRE: DAEMON
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable varnishncsa:
#
# varnishncsa_enable="YES"
#
# Configuration variables and their default values:
#
# varnishncsa_pidfile - full path to the PID file.
#	default: "/var/run/varnishncsa.pid"
#
# varnishncsa_file - full path to the log file.
#	default: "/var/log/varnishncsa.log"
#
# varnishncsa_flags - command line arguments.
#	default: "-D -P ${varnishncsa_pidfile} -a -c -w ${varnishncsa_file}"
#
# Add the following line to /etc/newsyslog.conf to rotate the log file
# once a day:
#
# /var/log/varnishncsa.log 640 7 * @T00 JB /var/run/varnishncsa.pid
#
# See varnishncsa(1) for a detailed overview of command-line options.
#

. /etc/rc.subr

name="varnishncsa"
rcvar=`set_rcvar`

command="%%PREFIX%%/bin/${name}"

# read configuration and set defaults
load_rc_config ${name}
: ${varnishncsa_enable:="NO"}
: ${varnishncsa_pidfile:="/var/run/${name}.pid"}
: ${varnishncsa_file:="/var/log/${name}.log"}
: ${varnishncsa_flags:="-P ${varnishncsa_pidfile} -D -a -c -w ${varnishncsa_file}"}

pidfile=${varnishncsa_pidfile}
run_rc_command "$1"
