# net/haproxy/files/openrc-haproxy.in
#!/sbin/openrc-run
#
# Add the following lines to /etc/rc.conf to enable haproxy:
#
# haproxy_pidfile (str):	default: /var/run/haproxy.pid
#				Set to the full path of the pid file
# haproxy_config (str):		default: %%PREFIX%%/etc/haproxy.conf
#				Set to the full path of the config file
# haproxy_flags (str):		default: Autogenerated using pidfile and config options
#				Set to override with your own options
# haproxy_profiles (str):	default: empty
# Set to space-separated list of profiles: for each profile separate haproxy
# process will be spawned, with haproxy-${profile}.conf config file.
# You can override default pidfile and config file for each profile with
# haproxy_${profile}_config and haproxy_${profile}_pidfile.

name="haproxy"
command="%%PREFIX%%/sbin/haproxy"

: ${haproxy_config:="%%PREFIX%%/etc/${name}.conf"}
pidfile=${haproxy_pidfile:-"/var/run/haproxy.pid"}
extracommands="reload configtest"

is_valid_profile() {
    local profile
    for profile in $haproxy_profiles; do
        if [ "$profile" = "$1" ]; then
            return 0
        fi
    done
    return 1
}

if [ -n "$2" ]; then
    profile=$2
    if ! is_valid_profile $profile; then
        eerror "$0: no such profile ($profile) defined in ${name}_profiles."
        exit 1
    fi
    eval haproxy_config="\${haproxy_${profile}_config:-%%PREFIX%%/etc/haproxy-${profile}.conf}"
    eval pidfile="\${haproxy_${profile}_pidfile:-/var/run/haproxy-${profile}.pid}"
else
    if [ "x${haproxy_profiles}" != "x" -a "x$1" != "x" ]; then
        for profile in ${haproxy_profiles}; do
            echo "===> ${name} profile: ${profile}"
            %%PREFIX%%/etc/init.d/haproxy $1 ${profile}
            retcode="$?"
            if [ ${retcode} -ne 0 ]; then
                failed="${profile} (${retcode}) ${failed:-}"
            else
                success="${profile} ${success:-}"
            fi
        done
        exit 0
    fi
fi

: ${haproxy_flags:="-q -f ${haproxy_config} -p ${pidfile}"}
required_files=$haproxy_config

configtest() {
	$command -c -f $haproxy_config
}

start_pre() {
	$command -q -c -f $haproxy_config
}

reload() {
	${command} -q -c -f ${haproxy_config}
	${command} ${haproxy_flags} -st $(cat ${pidfile})
}

depend() {
	keyword -shutdown
}
