# security/libreswan/pkg-plist
/etc/pam.d/pluto
etc/ipsec.conf
etc/ipsec.d/policies/block
etc/ipsec.d/policies/clear
etc/ipsec.d/policies/clear-or-private
etc/ipsec.d/policies/portexcludes.conf
etc/ipsec.d/policies/private
etc/ipsec.d/policies/private-or-clear
etc/ipsec.secrets
etc/logrotate.d/libreswan
etc/rc.d/pluto
libexec/ipsec/_import_crl
libexec/ipsec/_plutorun
libexec/ipsec/_realsetup
libexec/ipsec/_secretcensor
libexec/ipsec/_stackmanager
libexec/ipsec/_unbound-hook
libexec/ipsec/_updown
libexec/ipsec/_updown.bsdkame
libexec/ipsec/addconn
libexec/ipsec/algparse
libexec/ipsec/auto
libexec/ipsec/barf
libexec/ipsec/cavp
libexec/ipsec/dncheck
libexec/ipsec/ecdsasigkey
libexec/ipsec/enumcheck
libexec/ipsec/hunkcheck
libexec/ipsec/ipcheck
libexec/ipsec/jambufcheck
libexec/ipsec/keyidcheck
libexec/ipsec/letsencrypt
libexec/ipsec/look
libexec/ipsec/newhostkey
libexec/ipsec/pluto
libexec/ipsec/readwriteconf
libexec/ipsec/rsasigkey
libexec/ipsec/setup
libexec/ipsec/show
libexec/ipsec/showhostkey
libexec/ipsec/timecheck
libexec/ipsec/verify
libexec/ipsec/whack
sbin/ipsec
%%PORTDOCS%%%%DOCSDIR%%/ipsec.conf-sample
%%PORTDOCS%%%%DOCSDIR%%/ipsec.secrets-sample
share/examples/rc.d/pluto
share/man/man5/ipsec.conf.5.gz
share/man/man5/ipsec.secrets.5.gz
share/man/man8/ipsec.8.gz
share/man/man8/ipsec__import_crl.8.gz
share/man/man8/ipsec__plutorun.8.gz
share/man/man8/ipsec__realsetup.8.gz
share/man/man8/ipsec__secretcensor.8.gz
share/man/man8/ipsec__stackmanager.8.gz
share/man/man8/ipsec__unbound-hook.8.gz
share/man/man8/ipsec__updown.8.gz
share/man/man8/ipsec__updown.bsdkame.8.gz
share/man/man8/ipsec_addconn.8.gz
share/man/man8/ipsec_auto.8.gz
share/man/man8/ipsec_barf.8.gz
share/man/man8/ipsec_checknss.8.gz
share/man/man8/ipsec_ecdsasigkey.8.gz
share/man/man8/ipsec_import.8.gz
share/man/man8/ipsec_initnss.8.gz
share/man/man8/ipsec_letsencrypt.8.gz
share/man/man8/ipsec_look.8.gz
share/man/man8/ipsec_newhostkey.8.gz
share/man/man8/ipsec_pluto.8.gz
share/man/man8/ipsec_readwriteconf.8.gz
share/man/man8/ipsec_rsasigkey.8.gz
share/man/man8/ipsec_setup.8.gz
share/man/man8/ipsec_show.8.gz
share/man/man8/ipsec_showhostkey.8.gz
share/man/man8/ipsec_verify.8.gz
share/man/man8/ipsec_whack.8.gz
share/man/man8/pluto.8.gz
