# devel/linux-rl9-qt5-qtbase/pkg-plist.aarch64
etc/X11/xinit/xinitrc.d/10-qt5-check-opengl2.sh
etc/xdg/qtchooser/5-64.conf
usr/lib/.build-id/e2/9458e9ae35bc96c8fcca81ea13f44752de6459
usr/lib/.build-id/37/19857f9cc3831eb78b7c2ea0e0b103480bfd77
usr/lib/.build-id/1f/45b0a1728ccc0fea0b2bdeefdcd535f71b4d2d
usr/lib/.build-id/d1/6c5aa9cf460b1c2c9b736e69c1b46ef36126e3
usr/lib/.build-id/54/5443eb8101086947ace26159fa2d5fd322f898
usr/lib/.build-id/e5/08ba6a353099f2d4b51a91de7a4414ec746f47
usr/lib/.build-id/08/bfe1624c1761f11ebb0cd1908b4816f01ce337
usr/lib/.build-id/91/44821e9b9f3b712f5e8ead48d3de942f747f5d
usr/lib/.build-id/2c/b591ec4b33ea9c1f9d476433ccd4461ab0cfb3
usr/lib/.build-id/db/5cf2357a487edff865ba6effcb8a3c2a4e5161
usr/lib/.build-id/f9/ce6eccaae9d129edcb10d505137299d886f299
usr/lib/.build-id/a4/8d876bbc07ad2e56dfa089cd8f4c14a817b288
usr/lib/.build-id/4c/44f354eb7e7fa7a12c8673030418681949bb7f
usr/lib/.build-id/9e/e48fc1cb5de5d1a91fe4b5a45dfc84a031e729
usr/lib/.build-id/57/1938e23ead7a1f781024b82f7066c055105a2c
usr/lib/.build-id/ef/b66b6ec651fa81b932357da7bcf3a8f109b530
usr/lib/.build-id/1b/42356dbb348227804306a265d76ea78e578dc7
usr/lib/.build-id/a3/3a4780471f12e9ca2c011b5d21d392a5d06a2a
usr/lib/.build-id/b9/f6f3dd6bae1339e9ef73c5f46c42c20057de65
usr/lib/.build-id/b3/9689ee8f7fb3b44a96ec4ea07581b11b2fd767
usr/lib/.build-id/b3/3a80deb7fd14cb9b0d1f968b05420b540de6c6
usr/lib/.build-id/b3/244f8e551f43c2df27e2880ae328c46c48ade9
usr/lib/.build-id/f6/5d789f32deab00e8e6c6e8b8cb3c6e0c1872fe
usr/lib/.build-id/03/79a430fab8a459ab966a48480e495a65ffcfef
usr/lib/.build-id/22/f202afcdeaa128100aea8a528419a3acefc41b
usr/lib/.build-id/bd/5f098e308bd980ed580dac3ae554c4faf2d38d
usr/lib/.build-id/ec/a3702f3895c5d4d0071885a9e771256d05e17c
usr/lib/.build-id/ec/befb44f64069171f5944ee7bbbac96a4592bd1
usr/lib/.build-id/52/061e3cc36893a7d40a6f83dbc06ecaf0ed9b04
usr/lib/.build-id/a6/b23f54dec585dcc7aa6b0358b0ee62f30f6260
usr/lib/.build-id/a6/f19291df02b28af3ee509a7e353edc2d81c7b3
usr/lib/.build-id/e9/322b021690d52106af5aa1e3f25fbbd9830882
usr/lib/.build-id/36/22e298c6da2cd6654d0a82ca34018fb4234678
usr/lib/.build-id/36/81ab29d7f2f84c135aea630e669c06ba1d1a60
usr/lib/.build-id/93/ad1875c452b8f23068881f1238618bd13a40b8
usr/lib/.build-id/fe/8a514e5334b8f295a7315ccef68529141183d2
usr/lib/.build-id/21/dc03f8ed46ab38ad8bbade8f618b96042b53f8
usr/lib/.build-id/64/12731faf5c654634fc6f3e885701433eff02cc
usr/lib/.build-id/4b/9030d2235531c11ac64b2a7cef07468b41e913
usr/lib/.build-id/69/eebb97d304221b94a7d7f4ffdf4c377e6866c9
usr/lib/.build-id/26/d05e2fc86d14525bc756a44db6e8875f584acd
usr/lib/.build-id/9d/67ea603fcc16116aef45956325ab3d7e18f6f9
usr/lib/.build-id/c4/4a2fc657f60fb1e0ec02a03da1f949dcfedd27
usr/lib/.build-id/07/a407ec955cb71cd263d968cb183744db6e781e
usr/lib/.build-id/f2/8a2940dce5bf9176cd525f2f29c47b58836d63
usr/lib/.build-id/dc/f5281be1ae9dbd70d6bf4587656eb7f147063e
usr/lib/.build-id/94/a58c38d4cea14561710c2926c83487123e8ab1
usr/lib/.build-id/b7/ed81c28cdb72ede5142525af4b9aeb1e21cff9
usr/lib/rpm/macros.d/macros.qt5-qtbase
usr/lib64/cmake/Qt5Gui/Qt5Gui_QComposePlatformInputContextPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QEglFSEmulatorIntegrationPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QEglFSIntegrationPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QEglFSKmsEglDeviceIntegrationPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QEglFSKmsGbmIntegrationPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QEglFSX11IntegrationPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QEvdevKeyboardPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QEvdevMousePlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QEvdevTabletPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QEvdevTouchScreenPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QGifPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QGtk3ThemePlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QICOPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QIbusPlatformInputContextPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QJpegPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QLibInputPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QLinuxFbIntegrationPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QMinimalEglIntegrationPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QMinimalIntegrationPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QOffscreenIntegrationPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QTuioTouchPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QVncIntegrationPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QXcbEglIntegrationPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QXcbGlxIntegrationPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QXcbIntegrationPlugin.cmake
usr/lib64/cmake/Qt5Gui/Qt5Gui_QXdgDesktopPortalThemePlugin.cmake
usr/lib64/cmake/Qt5Network/Qt5Network_QConnmanEnginePlugin.cmake
usr/lib64/cmake/Qt5Network/Qt5Network_QGenericEnginePlugin.cmake
usr/lib64/cmake/Qt5Network/Qt5Network_QNetworkManagerEnginePlugin.cmake
usr/lib64/cmake/Qt5PrintSupport/Qt5PrintSupport_QCupsPrinterSupportPlugin.cmake
usr/lib64/cmake/Qt5Sql/Qt5Sql_QMYSQLDriverPlugin.cmake
usr/lib64/cmake/Qt5Sql/Qt5Sql_QODBCDriverPlugin.cmake
usr/lib64/cmake/Qt5Sql/Qt5Sql_QPSQLDriverPlugin.cmake
usr/lib64/cmake/Qt5Sql/Qt5Sql_QSQLiteDriverPlugin.cmake
usr/lib64/libQt5Concurrent.so.5
usr/lib64/libQt5Concurrent.so.5.15
usr/lib64/libQt5Concurrent.so.5.15.9
usr/lib64/libQt5Core.so.5
usr/lib64/libQt5Core.so.5.15
usr/lib64/libQt5Core.so.5.15.9
usr/lib64/libQt5DBus.so.5
usr/lib64/libQt5DBus.so.5.15
usr/lib64/libQt5DBus.so.5.15.9
usr/lib64/libQt5EglFSDeviceIntegration.so.5
usr/lib64/libQt5EglFSDeviceIntegration.so.5.15
usr/lib64/libQt5EglFSDeviceIntegration.so.5.15.9
usr/lib64/libQt5EglFsKmsSupport.so.5
usr/lib64/libQt5EglFsKmsSupport.so.5.15
usr/lib64/libQt5EglFsKmsSupport.so.5.15.9
usr/lib64/libQt5Gui.so.5
usr/lib64/libQt5Gui.so.5.15
usr/lib64/libQt5Gui.so.5.15.9
usr/lib64/libQt5Network.so.5
usr/lib64/libQt5Network.so.5.15
usr/lib64/libQt5Network.so.5.15.9
usr/lib64/libQt5OpenGL.so.5
usr/lib64/libQt5OpenGL.so.5.15
usr/lib64/libQt5OpenGL.so.5.15.9
usr/lib64/libQt5PrintSupport.so.5
usr/lib64/libQt5PrintSupport.so.5.15
usr/lib64/libQt5PrintSupport.so.5.15.9
usr/lib64/libQt5Sql.so.5
usr/lib64/libQt5Sql.so.5.15
usr/lib64/libQt5Sql.so.5.15.9
usr/lib64/libQt5Test.so.5
usr/lib64/libQt5Test.so.5.15
usr/lib64/libQt5Test.so.5.15.9
usr/lib64/libQt5Widgets.so.5
usr/lib64/libQt5Widgets.so.5.15
usr/lib64/libQt5Widgets.so.5.15.9
usr/lib64/libQt5XcbQpa.so.5
usr/lib64/libQt5XcbQpa.so.5.15
usr/lib64/libQt5XcbQpa.so.5.15.9
usr/lib64/libQt5Xml.so.5
usr/lib64/libQt5Xml.so.5.15
usr/lib64/libQt5Xml.so.5.15.9
usr/lib64/qt5/plugins/bearer/libqconnmanbearer.so
usr/lib64/qt5/plugins/bearer/libqgenericbearer.so
usr/lib64/qt5/plugins/bearer/libqnmbearer.so
usr/lib64/qt5/plugins/egldeviceintegrations/libqeglfs-emu-integration.so
usr/lib64/qt5/plugins/egldeviceintegrations/libqeglfs-kms-egldevice-integration.so
usr/lib64/qt5/plugins/egldeviceintegrations/libqeglfs-kms-integration.so
usr/lib64/qt5/plugins/egldeviceintegrations/libqeglfs-x11-integration.so
usr/lib64/qt5/plugins/generic/libqevdevkeyboardplugin.so
usr/lib64/qt5/plugins/generic/libqevdevmouseplugin.so
usr/lib64/qt5/plugins/generic/libqevdevtabletplugin.so
usr/lib64/qt5/plugins/generic/libqevdevtouchplugin.so
usr/lib64/qt5/plugins/generic/libqlibinputplugin.so
usr/lib64/qt5/plugins/generic/libqtuiotouchplugin.so
usr/lib64/qt5/plugins/imageformats/libqgif.so
usr/lib64/qt5/plugins/imageformats/libqico.so
usr/lib64/qt5/plugins/imageformats/libqjpeg.so
usr/lib64/qt5/plugins/platforminputcontexts/libcomposeplatforminputcontextplugin.so
usr/lib64/qt5/plugins/platforminputcontexts/libibusplatforminputcontextplugin.so
usr/lib64/qt5/plugins/platforms/libqeglfs.so
usr/lib64/qt5/plugins/platforms/libqlinuxfb.so
usr/lib64/qt5/plugins/platforms/libqminimal.so
usr/lib64/qt5/plugins/platforms/libqminimalegl.so
usr/lib64/qt5/plugins/platforms/libqoffscreen.so
usr/lib64/qt5/plugins/platforms/libqvnc.so
usr/lib64/qt5/plugins/platforms/libqxcb.so
usr/lib64/qt5/plugins/platformthemes/libqgtk3.so
usr/lib64/qt5/plugins/platformthemes/libqxdgdesktopportal.so
usr/lib64/qt5/plugins/printsupport/libcupsprintersupport.so
usr/lib64/qt5/plugins/sqldrivers/libqsqlite.so
usr/lib64/qt5/plugins/sqldrivers/libqsqlmysql.so
usr/lib64/qt5/plugins/sqldrivers/libqsqlodbc.so
usr/lib64/qt5/plugins/sqldrivers/libqsqlpsql.so
usr/lib64/qt5/plugins/xcbglintegrations/libqxcb-egl-integration.so
usr/lib64/qt5/plugins/xcbglintegrations/libqxcb-glx-integration.so
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qt3d.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtactiveqt.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtandroidextras.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtbase.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtcharts.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtconnectivity.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtdatavis3d.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtdeclarative.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtdoc.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtgamepad.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtgraphicaleffects.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtimageformats.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtlocation.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtlottie.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtmacextras.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtmultimedia.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtnetworkauth.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtpurchasing.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtquickcontrols.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtquickcontrols2.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtremoteobjects.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtscript.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtscxml.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtsensors.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtserialbus.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtserialport.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtspeech.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtsvg.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qttools.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtvirtualkeyboard.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtwayland.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtwebchannel.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtwebengine.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtwebsockets.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtwebview.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtwinextras.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtx11extras.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/config/exampleurl-qtxmlpatterns.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/compat.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/config.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/externalsites.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/externalsites/external-resources.qdoc
%%PORTDOCS%%%%DOCSDIR%%/global/externalsites/qt-webpages.qdoc
%%PORTDOCS%%%%DOCSDIR%%/global/externalsites/qtcreator.qdoc
%%PORTDOCS%%%%DOCSDIR%%/global/externalsites/rfc.qdoc
%%PORTDOCS%%%%DOCSDIR%%/global/fileextensions.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/html-config.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/html-footer-online.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/html-footer.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/html-header-offline.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/html-header-online.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/includes-online/search.qdoc
%%PORTDOCS%%%%DOCSDIR%%/global/includes/examples-run.qdocinc
%%PORTDOCS%%%%DOCSDIR%%/global/includes/module-use.qdocinc
%%PORTDOCS%%%%DOCSDIR%%/global/includes/standardpath/functiondocs.qdocinc
%%PORTDOCS%%%%DOCSDIR%%/global/macros.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/manifest-meta.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/qt-cpp-defines.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/qt-html-templates-offline-simple.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/qt-html-templates-offline.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/qt-html-templates-online.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/qt-module-defaults-offline.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/qt-module-defaults-online-commercial.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/qt-module-defaults-online.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/qt-module-defaults.qdocconf
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/Qt-dark_gradient.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/Qt-footer-bg.jpg
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/Qt-footer_shadow.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/Qt-gradient.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/Qt-header-bg.jpg
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/Qt-logo.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/arrow.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/arrow_bc.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/arrow_down.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/bg_l.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/bg_l_blank.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/bg_ll_blank.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/bg_r.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/bg_ul_blank.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/bgrContent.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/blu_dot.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/box_bg.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/breadcrumb.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/btn_next.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/btn_prev.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/bullet_dn.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/bullet_gt.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/bullet_sq.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/bullet_up.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/feedbackground.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/header_bg.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/home.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/horBar.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/ico_note.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/ico_note_attention.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/ico_out.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/logo.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/page.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/page_bg.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/spinner.gif
%%PORTDOCS%%%%DOCSDIR%%/global/template/images/sprites-combined.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/scripts/extras.js
%%PORTDOCS%%%%DOCSDIR%%/global/template/scripts/main.js
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/cookie-confirm.css
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/cookiebar-x.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/doc_search.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/gsc.css
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/icomoon.eot
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/icomoon.svg
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/icomoon.ttf
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/icomoon.woff
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/list_arrow.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/list_expand.png
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/offline-simple.css
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/offline.css
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/online.css
%%PORTDOCS%%%%DOCSDIR%%/global/template/style/theqtcompany.png
usr/share/licenses/qt5-qtbase/LICENSE.FDL
usr/share/licenses/qt5-qtbase/LICENSE.GPL2
usr/share/licenses/qt5-qtbase/LICENSE.GPL3
usr/share/licenses/qt5-qtbase/LICENSE.GPL3-EXCEPT
usr/share/licenses/qt5-qtbase/LICENSE.LGPL3
usr/share/licenses/qt5-qtbase/LICENSE.LGPLv3
%%DATADIR%%/qtlogging.ini
@dir etc/xdg/QtProject
@dir usr/lib64/cmake/Qt5
@dir usr/lib64/cmake/Qt5Concurrent
@dir usr/lib64/cmake/Qt5Core
@dir usr/lib64/cmake/Qt5DBus
@dir usr/lib64/cmake/Qt5OpenGL
@dir usr/lib64/cmake/Qt5Test
@dir usr/lib64/cmake/Qt5Widgets
@dir usr/lib64/cmake/Qt5Xml
@dir usr/lib64/qt5/imports
@dir usr/lib64/qt5/libexec
@dir usr/lib64/qt5/plugins/designer
@dir usr/lib64/qt5/plugins/iconengines
@dir usr/lib64/qt5/plugins/script
@dir usr/lib64/qt5/plugins/styles
@dir %%DATADIR%%/translations
