# mail/vbsfilter/files/vbsfilter.in
#!/bin/sh

# $FreeBSD$

# Start or stop vbsfilter

# PROVIDE: vbsfilter
# REQUIRE: DAEMON
# BEFORE: mail
# KEYWORD: shutdown
#

# Define these vbsfilter_* variables in one of these files:
#	/etc/rc.conf
#	/etc/rc.conf.local
#	/etc/rc.conf.d/vbsfilter
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#
# vbsfilter_flags Flags to vbsfilter

vbsfilter_enable=${vbsfilter_enable:-"NO"}      # Enable vbsfilter-milter
vbsfilter_socket=${vbsfilter_socket:-"/var/run/vbsfilter"} # Path to socket

. /etc/rc.subr

name="vbsfilter"
rcvar="vbsfilter_enable"
command="%%PREFIX%%/libexec/vbsfilter"
procname="${command}"
start_precmd="vbsfilter_prestart"
stop_postcmd="vbsfilter_poststop"

vbsfilter_prestart() {
	/bin/rm -f "${vbsfilter_socket}"
}

vbsfilter_poststop() {
	/bin/rm -f "${vbsfilter_socket}"
}

load_rc_config $name

command_args="-p local:${vbsfilter_socket}"

run_rc_command "$1"

# eof
