# net-mgmt/net-snmp4/pkg-plist
bin/encode_keychange
bin/snmpbulkwalk
bin/snmpcheck
bin/snmpdelta
bin/snmpget
bin/snmpgetnext
bin/snmpnetstat
bin/snmpset
bin/snmpstatus
bin/snmptable
bin/snmptest
bin/snmptranslate
bin/snmptrap
bin/snmpusm
bin/snmpwalk
bin/tkmib
etc/rc.d/snmpd.sh.sample
include/ucd-snmp/asn1.h
include/ucd-snmp/callback.h
include/ucd-snmp/default_store.h
include/ucd-snmp/int64.h
include/ucd-snmp/mib.h
include/ucd-snmp/parse.h
include/ucd-snmp/snmp.h
include/ucd-snmp/snmp_alarm.h
include/ucd-snmp/snmp_api.h
include/ucd-snmp/snmp_client.h
include/ucd-snmp/snmp_debug.h
include/ucd-snmp/snmp_impl.h
include/ucd-snmp/system.h
include/ucd-snmp/ucd-snmp-config.h
include/ucd-snmp/ucd-snmp-includes.h
include/ucd-snmp/version.h
@dirrm include/ucd-snmp
lib/libsnmp.so
lib/libsnmp.so.%%shlib%%
sbin/snmpd
sbin/snmptrapd
share/examples/ucd-snmp/README
share/examples/ucdDemoPublic.conf
@dirrm share/examples/ucd-snmp
share/snmp/mibs/DLMOD-MIB.txt
share/snmp/mibs/EtherLike-MIB.txt
share/snmp/mibs/HOST-RESOURCES-MIB.txt
share/snmp/mibs/IANAifType-MIB.txt
share/snmp/mibs/IF-MIB.txt
share/snmp/mibs/IP-MIB.txt
share/snmp/mibs/IPFWACC-MIB.txt
share/snmp/mibs/RFC1155-SMI.txt
share/snmp/mibs/RFC1213-MIB.txt
share/snmp/mibs/RFC1271-MIB.txt
share/snmp/mibs/SNMP-FRAMEWORK-MIB.txt
share/snmp/mibs/SNMP-MPD-MIB.txt
share/snmp/mibs/SNMP-NOTIFICATION-MIB.txt
share/snmp/mibs/SNMP-PROXY-MIB.txt
share/snmp/mibs/SNMP-TARGET-MIB.txt
share/snmp/mibs/SNMP-USER-BASED-SM-MIB.txt
share/snmp/mibs/SNMP-VIEW-BASED-ACM-MIB.txt
share/snmp/mibs/SNMPv2-CONF.txt
share/snmp/mibs/SNMPv2-M2M-MIB.txt
share/snmp/mibs/SNMPv2-MIB.txt
share/snmp/mibs/SNMPv2-PARTY-MIB.txt
share/snmp/mibs/SNMPv2-SMI.txt
share/snmp/mibs/SNMPv2-TC.txt
share/snmp/mibs/SNMPv2-TM.txt
share/snmp/mibs/TCP-MIB.txt
share/snmp/mibs/UCD-SNMP-MIB.txt
share/snmp/mibs/UDP-MIB.txt
@unexec rm -f %D/share/snmp/mibs/.index
@dirrm share/snmp/mibs
@dirrm share/snmp
@exec /usr/bin/env OBJFORMAT=%%PORTOBJFORMAT%% /sbin/ldconfig -m %D/lib
@unexec /usr/bin/env OBJFORMAT=%%PORTOBJFORMAT%% /sbin/ldconfig -R
