# net/freeradius/files/radiusd.sh.in
#!/bin/sh

# PROVIDE: radiusd
# REQUIRE: %%REQUIRE%%
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable radiusd:
#
# radiusd_enable="YES"
#

. /etc/rc.subr

name=radiusd
rcvar=`set_rcvar`
load_rc_config $name

command=%%PREFIX%%/sbin/radiusd
pidfile=/var/run/radiusd/radiusd.pid

# The radiusd_config parameter is now obsolete. Instead, derive the
# configuration directory name from radiusd_flags if possible, otherwise
# default to %%PREFIX%%/etc/raddb
required_dirs=`echo ${radiusd_flags} | sed -Ene 's:.*\-[^[:space:]d]*d[[:space:]]*([^[:space:]]+).*:\1:p'`
required_dirs=${required_dirs:-"%%PREFIX%%/etc/raddb"}

required_files="${required_dirs}/radiusd.conf"

radiusd_enable=${radiusd_enable-"NO"}

command_args="&"

run_rc_command "$1"
