# dns/dns_balance/files/openrc-dns_balance.in
#!/sbin/openrc-run

name=dns_balance
: ${dns_balance_logfile="/var/log/${name}.log"}
: ${dns_balance_pidfile="/var/log/${name}.pid"}

command="%%PREFIX%%/sbin/${name}"
command_args="-l ${dns_balance_logfile} -p ${dns_balance_pidfile}"
pidfile=${dns_balance_pidfile}
