# graphics/mesa-libs/pkg-plist
etc/libmap.d/mesa.conf
include/EGL/eglext_angle.h
include/EGL/eglmesaext.h
@comment include/GL/internal/dri_interface.h
include/gbm.h
@comment lib/dri/radeon_dri.so
@comment lib/dri/r300_dri.so
lib/libEGL_mesa.so
lib/libEGL_mesa.so.0
lib/libEGL_mesa.so.0.0.0
%%GLX%%lib/libGLX_mesa.so
%%GLX%%lib/libGLX_mesa.so.0
%%GLX%%lib/libGLX_mesa.so.0.0.0
lib/libgbm.so
lib/libgbm.so.1
lib/libgbm.so.1.0.0
lib/libglapi.so
lib/libglapi.so.0
lib/libglapi.so.0.0.0
@comment libdata/pkgconfig/dri.pc
libdata/pkgconfig/gbm.pc
@comment share/drirc.d/00-mesa-defaults.conf
share/glvnd/egl_vendor.d/50_mesa.json
