# graphics/mesa-libs/pkg-plist
include/EGL/egl.h
include/EGL/eglext.h
include/EGL/eglextchromium.h
include/EGL/eglmesaext.h
include/EGL/eglplatform.h
include/GL/gl.h
include/GL/gl_mangle.h
include/GL/glcorearb.h
include/GL/glext.h
include/GL/glx.h
include/GL/glx_mangle.h
include/GL/glxext.h
@comment include/GL/internal/dri_interface.h
include/GLES2/gl2.h
include/GLES2/gl2ext.h
include/GLES2/gl2platform.h
include/GLES3/gl3.h
include/GLES3/gl31.h
include/GLES3/gl32.h
include/GLES3/gl3ext.h
include/GLES3/gl3platform.h
include/KHR/khrplatform.h
include/gbm.h
@comment lib/dri/swrast_dri.so
lib/libEGL.so
lib/libEGL.so.1
lib/libEGL.so.1.0.0
lib/libGL.so
lib/libGL.so.1
lib/libGL.so.1.2.0
lib/libGLESv2.so
lib/libGLESv2.so.2
lib/libGLESv2.so.2.0.0
lib/libgbm.so
lib/libgbm.so.1
lib/libgbm.so.1.0.0
lib/libglapi.so
lib/libglapi.so.0
lib/libglapi.so.0.0.0
@comment libdata/pkgconfig/dri.pc
libdata/pkgconfig/egl.pc
libdata/pkgconfig/gbm.pc
libdata/pkgconfig/gl.pc
libdata/pkgconfig/glesv2.pc
@comment share/drirc.d/00-mesa-defaults.conf
