# net/freeradius/pkg-plist
bin/radclient
bin/radlast
bin/radrelay
bin/radtest
bin/radwho
bin/radzap
bin/rlm_dbm_cat
bin/rlm_dbm_parser
bin/rlm_ippool_tool
bin/smbencrypt
etc/raddb/acct_users
etc/raddb/attrs
etc/raddb/clients
etc/raddb/clients.conf
etc/raddb/dictionary
etc/raddb/experimental.conf
etc/raddb/hints
etc/raddb/huntgroups
etc/raddb/ldap.attrmap
etc/raddb/mssql.conf
etc/raddb/naslist
etc/raddb/naspasswd
etc/raddb/oraclesql.conf
etc/raddb/postgresql.conf
etc/raddb/preproxy_users
etc/raddb/proxy.conf
etc/raddb/radiusd.conf
etc/raddb/realms
etc/raddb/snmp.conf
etc/raddb/sql.conf
etc/raddb/users
etc/raddb/x99.conf
etc/raddb/x99passwd.sample
etc/rc.d/radiusd.sh
lib/libradius-0.9.2.la
lib/libradius-0.9.2.so
lib/libradius.a
lib/libradius.la
lib/libradius.so
lib/rlm_acct_unique-0.9.2.la
lib/rlm_acct_unique-0.9.2.so
lib/rlm_acct_unique.a
lib/rlm_acct_unique.la
lib/rlm_acct_unique.so
lib/rlm_always-0.9.2.la
lib/rlm_always-0.9.2.so
lib/rlm_always.a
lib/rlm_always.la
lib/rlm_always.so
lib/rlm_attr_filter-0.9.2.la
lib/rlm_attr_filter-0.9.2.so
lib/rlm_attr_filter.a
lib/rlm_attr_filter.la
lib/rlm_attr_filter.so
lib/rlm_attr_rewrite-0.9.2.la
lib/rlm_attr_rewrite-0.9.2.so
lib/rlm_attr_rewrite.a
lib/rlm_attr_rewrite.la
lib/rlm_attr_rewrite.so
lib/rlm_chap-0.9.2.la
lib/rlm_chap-0.9.2.so
lib/rlm_chap.a
lib/rlm_chap.la
lib/rlm_chap.so
lib/rlm_counter-0.9.2.la
lib/rlm_counter-0.9.2.so
lib/rlm_counter.a
lib/rlm_counter.la
lib/rlm_counter.so
lib/rlm_dbm-0.9.2.la
lib/rlm_dbm-0.9.2.so
lib/rlm_dbm.a
lib/rlm_dbm.la
lib/rlm_dbm.so
lib/rlm_detail-0.9.2.la
lib/rlm_detail-0.9.2.so
lib/rlm_detail.a
lib/rlm_detail.la
lib/rlm_detail.so
lib/rlm_digest-0.9.2.la
lib/rlm_digest-0.9.2.so
lib/rlm_digest.a
lib/rlm_digest.la
lib/rlm_digest.so
lib/rlm_eap-0.9.2.la
lib/rlm_eap-0.9.2.so
lib/rlm_eap.a
lib/rlm_eap.la
lib/rlm_eap.so
lib/rlm_eap_leap-0.9.2.la
lib/rlm_eap_leap-0.9.2.so
lib/rlm_eap_leap.a
lib/rlm_eap_leap.la
lib/rlm_eap_leap.so
lib/rlm_eap_md5-0.9.2.la
lib/rlm_eap_md5-0.9.2.so
lib/rlm_eap_md5.a
lib/rlm_eap_md5.la
lib/rlm_eap_md5.so
lib/rlm_eap_tls-0.9.2.la
lib/rlm_eap_tls-0.9.2.so
lib/rlm_eap_tls.a
lib/rlm_eap_tls.la
lib/rlm_eap_tls.so
lib/rlm_exec-0.9.2.la
lib/rlm_exec-0.9.2.so
lib/rlm_exec.a
lib/rlm_exec.la
lib/rlm_exec.so
lib/rlm_expr-0.9.2.la
lib/rlm_expr-0.9.2.so
lib/rlm_expr.a
lib/rlm_expr.la
lib/rlm_expr.so
lib/rlm_fastusers-0.9.2.la
lib/rlm_fastusers-0.9.2.so
lib/rlm_fastusers.a
lib/rlm_fastusers.la
lib/rlm_fastusers.so
lib/rlm_files-0.9.2.la
lib/rlm_files-0.9.2.so
lib/rlm_files.a
lib/rlm_files.la
lib/rlm_files.so
lib/rlm_ippool-0.9.2.la
lib/rlm_ippool-0.9.2.so
lib/rlm_ippool.a
lib/rlm_ippool.la
lib/rlm_ippool.so
%%LDAP%%lib/rlm_ldap-0.9.2.la
%%LDAP%%lib/rlm_ldap-0.9.2.so
%%LDAP%%lib/rlm_ldap.a
%%LDAP%%lib/rlm_ldap.la
%%LDAP%%lib/rlm_ldap.so
lib/rlm_mschap-0.9.2.la
lib/rlm_mschap-0.9.2.so
lib/rlm_mschap.a
lib/rlm_mschap.la
lib/rlm_mschap.so
lib/rlm_ns_mta_md5-0.9.2.la
lib/rlm_ns_mta_md5-0.9.2.so
lib/rlm_ns_mta_md5.a
lib/rlm_ns_mta_md5.la
lib/rlm_ns_mta_md5.so
lib/rlm_pam-0.9.2.la
lib/rlm_pam-0.9.2.so
lib/rlm_pam.a
lib/rlm_pam.la
lib/rlm_pam.so
lib/rlm_pap-0.9.2.la
lib/rlm_pap-0.9.2.so
lib/rlm_pap.a
lib/rlm_pap.la
lib/rlm_pap.so
lib/rlm_preprocess-0.9.2.la
lib/rlm_preprocess-0.9.2.so
lib/rlm_preprocess.a
lib/rlm_preprocess.la
lib/rlm_preprocess.so
lib/rlm_radutmp-0.9.2.la
lib/rlm_radutmp-0.9.2.so
lib/rlm_radutmp.a
lib/rlm_radutmp.la
lib/rlm_radutmp.so
lib/rlm_realm-0.9.2.la
lib/rlm_realm-0.9.2.so
lib/rlm_realm.a
lib/rlm_realm.la
lib/rlm_realm.so
lib/rlm_sql-0.9.2.la
lib/rlm_sql-0.9.2.so
lib/rlm_sql.a
lib/rlm_sql.la
lib/rlm_sql.so
%%MYSQL%%lib/rlm_sql_mysql.a
%%MYSQL%%lib/rlm_sql_mysql.la
%%MYSQL%%lib/rlm_sql_mysql.so
%%MYSQL%%lib/rlm_sql_mysql.so.0
lib/rlm_unix-0.9.2.la
lib/rlm_unix-0.9.2.so
lib/rlm_unix.a
lib/rlm_unix.la
lib/rlm_unix.so
%%TOKEN%%lib/rlm_x99_token-0.9.2.la
%%TOKEN%%lib/rlm_x99_token-0.9.2.so
%%TOKEN%%lib/rlm_x99_token.a
%%TOKEN%%lib/rlm_x99_token.la
%%TOKEN%%lib/rlm_x99_token.so
sbin/check-radiusd-config
sbin/checkrad
sbin/radiusd
sbin/radwatch
sbin/rc.radiusd
%%DOCSDIR%%/Autz-Type
%%DOCSDIR%%/DIFFS
%%DOCSDIR%%/MACOSX
%%DOCSDIR%%/OS2
%%DOCSDIR%%/RADIUS-LDAP.schema
%%DOCSDIR%%/RADIUS-LDAPv3.schema
%%DOCSDIR%%/RADIUS-SQL.schema
%%DOCSDIR%%/README
%%DOCSDIR%%/Simultaneous-Use
%%DOCSDIR%%/aaa.txt
%%DOCSDIR%%/ascend
%%DOCSDIR%%/bay
%%DOCSDIR%%/bugs
%%DOCSDIR%%/cisco
%%DOCSDIR%%/coding-methods.txt
%%DOCSDIR%%/configurable_failover
%%DOCSDIR%%/duplicate-users
%%DOCSDIR%%/module_interface
%%DOCSDIR%%/performance-testing
%%DOCSDIR%%/processing_users_file
%%DOCSDIR%%/proxy
%%DOCSDIR%%/radrelay
%%DOCSDIR%%/rfc/attributes.html
%%DOCSDIR%%/rfc/draft-sterman-aaa-sip-00.txt
%%DOCSDIR%%/rfc/leap.txt
%%DOCSDIR%%/rfc/rfc1157.txt
%%DOCSDIR%%/rfc/rfc1227.txt
%%DOCSDIR%%/rfc/rfc1448.txt
%%DOCSDIR%%/rfc/rfc1901.txt
%%DOCSDIR%%/rfc/rfc1905.txt
%%DOCSDIR%%/rfc/rfc2058.txt
%%DOCSDIR%%/rfc/rfc2059.txt
%%DOCSDIR%%/rfc/rfc2138.txt
%%DOCSDIR%%/rfc/rfc2139.txt
%%DOCSDIR%%/rfc/rfc2243.txt
%%DOCSDIR%%/rfc/rfc2289.txt
%%DOCSDIR%%/rfc/rfc2433.txt
%%DOCSDIR%%/rfc/rfc2548.txt
%%DOCSDIR%%/rfc/rfc2618.txt
%%DOCSDIR%%/rfc/rfc2619.txt
%%DOCSDIR%%/rfc/rfc2620.txt
%%DOCSDIR%%/rfc/rfc2621.txt
%%DOCSDIR%%/rfc/rfc2716.txt
%%DOCSDIR%%/rfc/rfc2759.txt
%%DOCSDIR%%/rfc/rfc2809.txt
%%DOCSDIR%%/rfc/rfc2865.txt
%%DOCSDIR%%/rfc/rfc2866.txt
%%DOCSDIR%%/rfc/rfc2867.txt
%%DOCSDIR%%/rfc/rfc2868.txt
%%DOCSDIR%%/rfc/rfc2869.txt
%%DOCSDIR%%/rfc/rfc2882.txt
%%DOCSDIR%%/rfc/rfc2924.txt
%%DOCSDIR%%/rfc/rfc3162.txt
%%DOCSDIR%%/rlm_acct_unique
%%DOCSDIR%%/rlm_always
%%DOCSDIR%%/rlm_attr_filter
%%DOCSDIR%%/rlm_dbm
%%DOCSDIR%%/rlm_detail
%%DOCSDIR%%/rlm_digest
%%DOCSDIR%%/rlm_eap
%%DOCSDIR%%/rlm_expr
%%DOCSDIR%%/rlm_fastusers
%%DOCSDIR%%/rlm_krb5
%%DOCSDIR%%/rlm_ldap
%%DOCSDIR%%/rlm_mschap
%%DOCSDIR%%/rlm_pam
%%DOCSDIR%%/rlm_passwd
%%DOCSDIR%%/rlm_python
%%DOCSDIR%%/rlm_sql
%%DOCSDIR%%/rlm_sqlcounter
%%DOCSDIR%%/rlm_unix
%%DOCSDIR%%/rlm_x99_token
%%DOCSDIR%%/supervise-radiusd.txt
%%DOCSDIR%%/tuning_guide
%%DOCSDIR%%/variables.txt
%%DATADIR%%/dictionary
%%DATADIR%%/dictionary.acc
%%DATADIR%%/dictionary.alcatel
%%DATADIR%%/dictionary.alteon
%%DATADIR%%/dictionary.altiga
%%DATADIR%%/dictionary.aptis
%%DATADIR%%/dictionary.ascend
%%DATADIR%%/dictionary.bay
%%DATADIR%%/dictionary.bintec
%%DATADIR%%/dictionary.cabletron
%%DATADIR%%/dictionary.cisco
%%DATADIR%%/dictionary.cisco.bbsm
%%DATADIR%%/dictionary.cisco.vpn3000
%%DATADIR%%/dictionary.cisco.vpn5000
%%DATADIR%%/dictionary.colubris
%%DATADIR%%/dictionary.columbia_university
%%DATADIR%%/dictionary.compat
%%DATADIR%%/dictionary.erx
%%DATADIR%%/dictionary.extreme
%%DATADIR%%/dictionary.foundry
%%DATADIR%%/dictionary.freeradius
%%DATADIR%%/dictionary.itk
%%DATADIR%%/dictionary.juniper
%%DATADIR%%/dictionary.livingston
%%DATADIR%%/dictionary.microsoft
%%DATADIR%%/dictionary.netscreen
%%DATADIR%%/dictionary.nokia
%%DATADIR%%/dictionary.nomadix
%%DATADIR%%/dictionary.quintum
%%DATADIR%%/dictionary.redback
%%DATADIR%%/dictionary.redcreek
%%DATADIR%%/dictionary.shasta
%%DATADIR%%/dictionary.shiva
%%DATADIR%%/dictionary.springtide
%%DATADIR%%/dictionary.telebit
%%DATADIR%%/dictionary.tunnel
%%DATADIR%%/dictionary.usr
%%DATADIR%%/dictionary.versanet
%%DATADIR%%/dictionary.wispr
%%DATADIR%%/dictionary.xedia
@dirrm %%DATADIR%%
@dirrm %%DOCSDIR%%/rfc
@dirrm %%DOCSDIR%%
@dirrm etc/raddb
@exec mkdir -p /var/log/raddb
@exec chmod -R og= /var/log/raddb
@unexec rm -fr /var/log/raddb
@exec mkdir -p /var/run/radiusd
@unexec rm -fr /var/run/radiusd
