# mail/exim/pkg-plist
etc/exim/configure.default
etc/rc.d/exim.sh
info/exim.info
info/exim_filter.info
sbin/exim
sbin/exim-%%EXIM_VERSION%%-1
sbin/exim_checkaccess
sbin/exim_dumpdb
sbin/exim_fixdb
sbin/exim_lock
sbin/exim_tidydb
sbin/exinext
sbin/exiwhat
sbin/exim_dbmbuild
sbin/exicyclog
sbin/exigrep
sbin/eximstats
sbin/exiqgrep
sbin/exiqsumm
%%PORTDOCS%%share/doc/exim/POST-INSTALL-NOTES
%%PORTDOCS%%share/doc/exim/Exim3.upgrade
%%PORTDOCS%%share/doc/exim/Exim4.upgrade
%%PORTDOCS%%share/doc/exim/OptionLists.txt
%%PORTDOCS%%share/doc/exim/README
%%PORTDOCS%%share/doc/exim/dbm.discuss.txt
%%PORTDOCS%%share/doc/exim/filter.txt
%%PORTDOCS%%share/doc/exim/pcre.txt
%%PORTDOCS%%share/doc/exim/pcretest.txt
%%PORTDOCS%%share/doc/exim/spec.txt
%%PORTDOCS%%%%EXISCAN%%share/doc/exim/exiscan-readme.txt
@unexec install-info --delete %D/info/exim.info %D/info/dir
@unexec install-info --delete %D/info/exim_filter.info %D/info/dir
@exec install-info --section="Exim" --entry "* Exim Specification: (exim).       Exim manual" %D/info/exim.info %D/info/dir
@exec install-info --section="Exim" --entry "* Exim Filtering: (exim_filter).    Filtering mail with Exim" %D/info/exim_filter.info %D/info/dir
@dirrm etc/exim
@dirrm share/doc/exim
