# mail/exim/pkg-plist
etc/exim/configure.default
etc/rc.d/exim.sh
sbin/exim
sbin/exim-%%EXIM_VERSION%%-1
sbin/exim_checkaccess
sbin/exim_dumpdb
sbin/exim_fixdb
sbin/exim_lock
sbin/exim_tidydb
sbin/exinext
sbin/exiwhat
sbin/exim_dbmbuild
sbin/exicyclog
sbin/exigrep
sbin/eximstats
sbin/exiqgrep
sbin/exiqsumm
%%PORTDOCS%%%%DOCSDIR%%/POST-INSTALL-NOTES
%%PORTDOCS%%%%DOCSDIR%%/Exim3.upgrade
%%PORTDOCS%%%%DOCSDIR%%/Exim4.upgrade
%%PORTDOCS%%%%DOCSDIR%%/OptionLists.txt
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/dbm.discuss.txt
%%PORTDOCS%%%%DOCSDIR%%/filter.txt
%%PORTDOCS%%%%DOCSDIR%%/pcrepattern.txt
%%PORTDOCS%%%%DOCSDIR%%/pcretest.txt
%%PORTDOCS%%%%DOCSDIR%%/spec.txt
%%PORTDOCS%%%%EXISCAN_ACL%%%%DOCSDIR%%/exiscan-acl-examples.txt
%%PORTDOCS%%%%EXISCAN_ACL%%%%DOCSDIR%%/exiscan-acl-spec.txt
@dirrm etc/exim
%%PORTDOCS%%@dirrm %%DOCSDIR%%
