# mail/exim/pkg-plist
etc/exim/configure.default
etc/rc.d/exim.sh
info/exim_overview.info
info/exim.info
info/exim_filter.info
sbin/exim
sbin/exim_dumpdb
sbin/exim_fixdb
sbin/exim_lock
sbin/exim_tidydb
sbin/exinext
sbin/exiwhat
sbin/exim_dbmbuild
sbin/exicyclog
sbin/exigrep
sbin/eximstats
sbin/exiqsumm
%%PORTDOCS%%share/doc/exim/POST-INSTALL-NOTES
%%PORTDOCS%%share/doc/exim/filter.txt
%%PORTDOCS%%share/doc/exim/oview.txt
%%PORTDOCS%%share/doc/exim/spec.txt
@unexec install-info --delete %D/info/exim_overview.info %D/info/dir
@unexec install-info --delete %D/info/exim.info %D/info/dir
@unexec install-info --delete %D/info/exim_filter.info %D/info/dir
@exec install-info --section="Exim" --entry "* Exim Overview: (exim_overview).   Overview of the Exim system" %D/info/exim_overview.info %D/info/dir
@exec install-info --section="Exim" --entry "* Exim Specification: (exim).       Exim manual" %D/info/exim.info %D/info/dir
@exec install-info --section="Exim" --entry "* Exim Filtering: (exim_filter).    Filtering mail with Exim" %D/info/exim_filter.info %D/info/dir
@dirrm etc/exim
@dirrm share/doc/exim
