# security/shibboleth2-sp/files/shibboleth-sp.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: shibboleth_sp
# REQUIRE: DAEMON
#

. /etc/rc.subr

name="shibboleth_sp"
rcvar=${name}_enable

[ -z "$shibboleth_sp_enable" ] && shibboleth_sp_enable="NO"
pidfile="/var/run/${name}.pid"
command="%%PREFIX%%/sbin/shibd"
start_cmd="${command} -f -p ${pidfile} &"

load_rc_config $name
run_rc_command "$1"
