# security/shibboleth2-sp/files/shibboleth-sp.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: shibboleth_sp
# REQUIRE: DAEMON
#

. /etc/rc.subr

name="shibboleth_sp"
rcvar=${name}_enable

command=${shibboleth_sp_program:-%%PREFIX%%/sbin/shibd}
pidfile="${shibboleth_sp_pidfile:-/var/run/${name}.pid}"
command_args="-f -p ${pidfile}"

load_rc_config $name
run_rc_command "$1"
