# security/expiretable/pkg-plist
sbin/expiretable
@unexec %%RC_DIR%%/etc/rc.d/expiretable%%RC_SUFX%% stop 2>&1 >/dev/null || true
