# databases/courier-authlib-userdb/pkg-plist
@unexec [ ! -f /var/run/authdaemond/pid ] || %D/etc/rc.d/courier-authdaemond.sh stop || true
bin/courierauthconfig
etc/rc.d/courier-authdaemond.sh
%%MYSQLFLAG%%etc/authlib/authmysqlrc.dist
%%PGSQLFLAG%%etc/authlib/authpgsqlrc.dist
%%LDAPFLAG%%etc/authlib/authldaprc.dist
@unexec if cmp -s %D/etc/authlib/authdaemonrc.dist %D/etc/authlib/authdaemonrc; then rm -f %D/etc/authlib/authdaemonrc; fi
etc/authlib/authdaemonrc.dist
include/courier_auth_config.h
include/courierauth.h
include/courierauthsasl.h
include/courierauthsaslclient.h
include/courierauthdebug.h
lib/courier-authlib/libcourierauth.so.0
lib/courier-authlib/libcourierauth.so
lib/courier-authlib/libcourierauth.a
lib/courier-authlib/libcourierauthsasl.so.0
lib/courier-authlib/libcourierauthsasl.so
lib/courier-authlib/libcourierauthsasl.a
lib/courier-authlib/libcourierauthsaslclient.so.0
lib/courier-authlib/libcourierauthsaslclient.so
lib/courier-authlib/libcourierauthsaslclient.a
lib/courier-authlib/libcourierauthcommon.so.0
lib/courier-authlib/libcourierauthcommon.so
lib/courier-authlib/libcourierauthcommon.a
%%AUTHUSERDBFLAG%%lib/courier-authlib/libauthuserdb.so.0
%%AUTHUSERDBFLAG%%lib/courier-authlib/libauthuserdb.so
%%PAMFLAG%%lib/courier-authlib/libauthpam.so.0
%%PAMFLAG%%lib/courier-authlib/libauthpam.so
%%MYSQLFLAG%%lib/courier-authlib/libauthmysql.so.0
%%MYSQLFLAG%%lib/courier-authlib/libauthmysql.so
%%PGSQLFLAG%%lib/courier-authlib/libauthpgsql.so.0
%%PGSQLFLAG%%lib/courier-authlib/libauthpgsql.so
%%LDAPFLAG%%lib/courier-authlib/libauthldap.so.0
%%LDAPFLAG%%lib/courier-authlib/libauthldap.so
%%VPOPMAILFLAG%%lib/courier-authlib/libauthvchkpw.so.0
%%VPOPMAILFLAG%%lib/courier-authlib/libauthvchkpw.so
libexec/courier-authlib/authsystem.passwd
libexec/courier-authlib/authdaemond
%%AUTHUSERDBFLAG%%libexec/courier-authlib/makedatprog
sbin/authtest
sbin/authdaemond
sbin/authenumerate
sbin/courierlogger
%%AUTHUSERDBFLAG%%sbin/pw2userdb
%%AUTHUSERDBFLAG%%sbin/makeuserdb
%%AUTHUSERDBFLAG%%sbin/userdb
%%AUTHUSERDBFLAG%%sbin/userdb-test-cram-md5
%%AUTHUSERDBFLAG%%sbin/userdbpw
%%AUTHUSERDBFLAG%%sbin/vchkpw2userdb
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README.authdebug.html
%%PORTDOCS%%%%DOCSDIR%%/README.authmysql.html
%%PORTDOCS%%%%DOCSDIR%%/README.authmysql.myownquery
%%PORTDOCS%%%%DOCSDIR%%/README.authpostgres.html
%%PORTDOCS%%%%DOCSDIR%%/README.ldap
%%PORTDOCS%%%%DOCSDIR%%/README_authlib.html
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@exec mkdir -m 0750 /var/run/authdaemond 2>/dev/null || true
@exec chown %%MAILUSER%%:%%MAILGROUP%% /var/run/authdaemond 2>/dev/null || true
@unexec rmdir /var/run/authdaemond 2>/dev/null || true
@unexec rmdir %D/etc/authlib 2>/dev/null || true
@dirrm lib/courier-authlib
@dirrm libexec/courier-authlib
