# misc/zaptel/pkg-plist
bin/fxo_test
bin/ztcfg
bin/ztmonitor
bin/zttest
bin/zttool
etc/rc.d/zaptel.sh.sample
include/tonezone.h
include/zaptel.h
lib/libtonezone.a
lib/zaptel/linker.hints
lib/zaptel/qozap.ko
lib/zaptel/wcfxo.ko
lib/zaptel/wcfxs.ko
lib/zaptel/wct1xxp.ko
lib/zaptel/wct4xxp.ko
lib/zaptel/wcte11xp.ko
lib/zaptel/zaptel.ko
lib/zaptel/ztdummy.ko
@dirrm lib/zaptel
