# www/shellinabox/pkg-plist
bin/shellinaboxd
etc/rc.d/shellinaboxd
man/man1/shellinaboxd.1.gz
@mode 700
@owner shellinabox
@group shellinabox
@exec mkdir -p %D/etc/shellinabox
@dirrmtry etc/shellinabox
@mode
