# sysutils/powerman/pkg-plist
bin/pm
bin/powerman
etc/powerman/apc.dev
etc/powerman/apcnew.dev
etc/powerman/baytech-rpc28-nc.dev
etc/powerman/baytech-rpc3-nc.dev
etc/powerman/baytech.dev
etc/powerman/cyclades-pm8.dev
etc/powerman/ibmh8.dev
etc/powerman/icebox.dev
etc/powerman/icebox3.dev
etc/powerman/ipmi.dev
etc/powerman/ipmipower.dev
etc/powerman/phantom.dev
etc/powerman/phantom4.dev
etc/powerman/vpc.dev
etc/powerman/wti-rps10.dev
etc/powerman/wti.dev
etc/rc.d/powermand.sh
sbin/powermand
@dirrm etc/powerman
