# security/drweb/pkg-plist
@unexec if [ -f /var/drweb/run/drwebd.pid ] ; then %D/etc/rc.d/drwebd stop ; fi
drweb/agent/daemon.amc
drweb/drweb
drweb/drwebd
drweb/drwebdc
drweb/getkey.HOWTO
drweb/getkey.rus.HOWTO
drweb/getkey.rus.HOWTO.utf8
drweb/lib/drweb32.dll
drweb/lib/ru_daemon.dwl
drweb/lib/ru_scanner.dwl
drweb/read_signed
drweb/update/update.pl
@unexec if [ -f %D/etc/drweb/drweb32.ini ] && cmp -s %D/etc/drweb/drweb32.ini %D/etc/drweb/drweb32.ini-distr; then rm -f %D/etc/drweb/drweb32.ini ; fi
etc/drweb/drweb32.ini-distr
@exec if [ ! -f %B/drweb32.ini ] ; then cp -p %B/%F %B/drweb32.ini ; fi
@unexec if [ -f %D/etc/drweb/email.ini ] && cmp -s %D/etc/drweb/email.ini %D/etc/drweb/email.ini-distr; then rm -f %D/etc/drweb/email.ini ; fi
etc/drweb/email.ini-distr
@exec if [ ! -f %B/email.ini ] ; then cp -p %B/%F %B/email.ini ; fi
etc/drweb/monitor/daemon.mmc
etc/rc.d/drwebd
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog.rus
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/FAQ
%%PORTDOCS%%%%DOCSDIR%%/FAQ.rus
%%PORTDOCS%%%%DOCSDIR%%/FAQ.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/daemon/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/daemon/ChangeLog.rus
%%PORTDOCS%%%%DOCSDIR%%/daemon/ChangeLog.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/daemon/readme.daemon
%%PORTDOCS%%%%DOCSDIR%%/daemon/readme.daemon.rus
%%PORTDOCS%%%%DOCSDIR%%/daemon/readme.daemon.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/readme.eicar
%%PORTDOCS%%%%DOCSDIR%%/readme.eicar.rus
%%PORTDOCS%%%%DOCSDIR%%/readme.eicar.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/readme.license
%%PORTDOCS%%%%DOCSDIR%%/readme.license.rus
%%PORTDOCS%%%%DOCSDIR%%/readme.license.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/scanner/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/scanner/ChangeLog.rus
%%PORTDOCS%%%%DOCSDIR%%/scanner/ChangeLog.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/scanner/readme.scanner
%%PORTDOCS%%%%DOCSDIR%%/scanner/readme.scanner.rus
%%PORTDOCS%%%%DOCSDIR%%/scanner/readme.scanner.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/update/readme.update
%%PORTDOCS%%%%DOCSDIR%%/update/readme.update.rus
%%PORTDOCS%%%%DOCSDIR%%/update/readme.update.rus.utf8
%%PORTDOCS%%@dirrm %%DOCSDIR%%/update
%%PORTDOCS%%@dirrm %%DOCSDIR%%/scanner
%%PORTDOCS%%@dirrm %%DOCSDIR%%/daemon
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm etc/drweb/monitor
@dirrm etc/drweb
@dirrm drweb/update
@dirrm drweb/lib
@dirrm drweb/agent
@dirrm drweb
@cwd /var
drweb/bases/drw44400.vdb
drweb/bases/drw44401.vdb
drweb/bases/drw44402.vdb
drweb/bases/drw44403.vdb
drweb/bases/drw44404.vdb
drweb/bases/drw44405.vdb
drweb/bases/drw44406.vdb
drweb/bases/drw44407.vdb
drweb/bases/drw44408.vdb
drweb/bases/drw44409.vdb
drweb/bases/drw44410.vdb
drweb/bases/drw44411.vdb
drweb/bases/drw44412.vdb
drweb/bases/drw44413.vdb
drweb/bases/drw44414.vdb
drweb/bases/drw44415.vdb
drweb/bases/drw44416.vdb
drweb/bases/drw44417.vdb
drweb/bases/drw44418.vdb
drweb/bases/drw44419.vdb
drweb/bases/drw44420.vdb
drweb/bases/drw44421.vdb
drweb/bases/drw44422.vdb
drweb/bases/drw44423.vdb
drweb/bases/drw44424.vdb
drweb/bases/drw44425.vdb
drweb/bases/drw44426.vdb
drweb/bases/drw44427.vdb
drweb/bases/drw44428.vdb
drweb/bases/drw44429.vdb
drweb/bases/drw44430.vdb
drweb/bases/drw44431.vdb
drweb/bases/drw44432.vdb
drweb/bases/drw44433.vdb
drweb/bases/drw44434.vdb
drweb/bases/drw44435.vdb
drweb/bases/drw44436.vdb
drweb/bases/drw44437.vdb
drweb/bases/drw44438.vdb
drweb/bases/drw44439.vdb
drweb/bases/drw44440.vdb
drweb/bases/drw44441.vdb
drweb/bases/drw44442.vdb
drweb/bases/drw44443.vdb
drweb/bases/drw44444.vdb
drweb/bases/drw44445.vdb
drweb/bases/drw44446.vdb
drweb/bases/drw44447.vdb
drweb/bases/drw44448.vdb
drweb/bases/drw44449.vdb
drweb/bases/drw44450.vdb
drweb/bases/drw44451.vdb
drweb/bases/drw44452.vdb
drweb/bases/drw44453.vdb
drweb/bases/drw44454.vdb
drweb/bases/drw44455.vdb
drweb/bases/drw44456.vdb
drweb/bases/drw44457.vdb
drweb/bases/drw44458.vdb
drweb/bases/drw44459.vdb
drweb/bases/drw44460.vdb
drweb/bases/drw44461.vdb
drweb/bases/drw44462.vdb
drweb/bases/drw44463.vdb
drweb/bases/drw44464.vdb
drweb/bases/drwebase.vdb
drweb/bases/drwnasty.vdb
drweb/bases/drwrisky.vdb
drweb/bases/drwtoday.vdb
drweb/bases/dwn44401.vdb
drweb/bases/dwn44402.vdb
drweb/bases/dwn44403.vdb
drweb/bases/dwn44404.vdb
drweb/bases/dwn44405.vdb
drweb/bases/dwn44406.vdb
drweb/bases/dwntoday.vdb
drweb/bases/dwr44401.vdb
drweb/bases/dwrtoday.vdb
drweb/bases/update.drl
@dirrm drweb/bases
@dirrm drweb/cache
@dirrm drweb/infected
@dirrm drweb/log
@dirrm drweb/run
@dirrm drweb/spool
@dirrm drweb/updates
@dirrm drweb
