# security/drweb/pkg-plist
bin/drweb
drweb/bases/drw43201.vdb
drweb/bases/drwebase.vdb
drweb/doc/ChangeLog
drweb/doc/FAQ
drweb/doc/FAQ.rus
drweb/doc/readme.daemon
drweb/doc/readme.daemon.rus
drweb/doc/readme.license
drweb/doc/readme.license.rus
drweb/doc/readme.scanner
drweb/doc/readme.scanner.rus
drweb/drweb
@unexec if [ -f %D/drweb/drweb32.ini ] && cmp -s %D/drweb/drweb32.ini %D/drweb/drweb32.ini-distr; then rm -f %D/drweb/drweb32.ini; fi
drweb/drweb32.ini-distr
@exec if [ ! -f %D/drweb/drweb32.ini ] ; then cp -p %D/%F %B/drweb32.ini; fi
drweb/drwebd
drweb/drwebdc
@unexec if [ -f %D/drweb/email.ini ] && cmp -s %D/drweb/email.ini %D/drweb/email.ini-distr; then rm -f %D/drweb/email.ini; fi
drweb/email.ini-distr
@exec if [ ! -f %D/drweb/email.ini ] ; then cp -p %D/%F %B/email.ini; fi
drweb/getkey.HOWTO
drweb/getkey.rus.HOWTO
drweb/lib/drweb32.dll
drweb/lib/russian.dwl
drweb/update/readme.update
drweb/update/readme.update.rus
drweb/update/update.pl
etc/rc.d/00.drwebd.sh
@dirrm drweb/updates
@dirrm drweb/update
@dirrm drweb/run
@dirrm drweb/lib
@dirrm drweb/infected
@dirrm drweb/doc
@dirrm drweb/bases
@dirrm drweb
