# security/drweb/pkg-plist
@unexec if [ -f /var/drweb/run/drwebd.pid ] ; then %D/etc/rc.d/drwebd stop ; fi
drweb/drweb
drweb/drwebd
drweb/drwebdc
drweb/getkey.HOWTO
drweb/getkey.rus.HOWTO
drweb/lib/drweb32.dll
drweb/lib/ru_daemon.dwl
drweb/lib/ru_scanner.dwl
drweb/read_signed
drweb/update/update.pl
@unexec if [ -f %D/etc/drweb/drweb32.ini ] && cmp -s %D/etc/drweb/drweb32.ini %D/etc/drweb/drweb32.ini-distr; then rm -f %D/etc/drweb/drweb32.ini ; fi
etc/drweb/drweb32.ini-distr
@exec if [ ! -f %B/drweb32.ini ] ; then cp -p %B/%F %B/drweb32.ini ; fi
@unexec if [ -f %D/etc/drweb/email.ini ] && cmp -s %D/etc/drweb/email.ini %D/etc/drweb/email.ini-distr; then rm -f %D/etc/drweb/email.ini ; fi
etc/drweb/email.ini-distr
@exec if [ ! -f %B/email.ini ] ; then cp -p %B/%F %B/email.ini ; fi
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog.rus
%%PORTDOCS%%%%DOCSDIR%%/FAQ
%%PORTDOCS%%%%DOCSDIR%%/FAQ.rus
%%PORTDOCS%%%%DOCSDIR%%/daemon/readme.daemon
%%PORTDOCS%%%%DOCSDIR%%/daemon/readme.daemon.rus
%%PORTDOCS%%%%DOCSDIR%%/drweb.newsyslog.add
%%PORTDOCS%%%%DOCSDIR%%/readme.eicar
%%PORTDOCS%%%%DOCSDIR%%/readme.eicar.rus
%%PORTDOCS%%%%DOCSDIR%%/readme.license
%%PORTDOCS%%%%DOCSDIR%%/readme.license.rus
%%PORTDOCS%%%%DOCSDIR%%/scanner/readme.scanner
%%PORTDOCS%%%%DOCSDIR%%/scanner/readme.scanner.rus
%%PORTDOCS%%%%DOCSDIR%%/update/readme.update
%%PORTDOCS%%%%DOCSDIR%%/update/readme.update.rus
%%PORTDOCS%%@dirrm %%DOCSDIR%%/update
%%PORTDOCS%%@dirrm %%DOCSDIR%%/scanner
%%PORTDOCS%%@dirrm %%DOCSDIR%%/daemon
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm etc/drweb
@dirrm drweb/update
@dirrm drweb/lib
@dirrm drweb
@cwd /var
drweb/bases/drw43300.vdb
drweb/bases/drw43301.vdb
drweb/bases/drw43302.vdb
drweb/bases/drw43303.vdb
drweb/bases/drw43304.vdb
drweb/bases/drw43305.vdb
drweb/bases/drw43306.vdb
drweb/bases/drw43307.vdb
drweb/bases/drw43308.vdb
drweb/bases/drw43309.vdb
drweb/bases/drw43310.vdb
drweb/bases/drw43311.vdb
drweb/bases/drw43312.vdb
drweb/bases/drw43313.vdb
drweb/bases/drw43314.vdb
drweb/bases/drw43315.vdb
drweb/bases/drw43316.vdb
drweb/bases/drw43317.vdb
drweb/bases/drw43318.vdb
drweb/bases/drw43319.vdb
drweb/bases/drw43320.vdb
drweb/bases/drw43321.vdb
drweb/bases/drw43322.vdb
drweb/bases/drw43323.vdb
drweb/bases/drw43324.vdb
%%RU%%drweb/bases/drw43325.vdb
%%RU%%drweb/bases/drw43326.vdb
%%RU%%drweb/bases/drw43327.vdb
%%RU%%drweb/bases/drw43328.vdb
%%RU%%drweb/bases/drw43329.vdb
%%RU%%drweb/bases/drw43330.vdb
%%RU%%drweb/bases/drw43331.vdb
%%RU%%drweb/bases/drw43332.vdb
%%RU%%drweb/bases/drw43333.vdb
%%RU%%drweb/bases/drw43334.vdb
%%RU%%drweb/bases/drw43335.vdb
%%RU%%drweb/bases/drw43336.vdb
%%RU%%drweb/bases/drw43337.vdb
%%RU%%drweb/bases/drw43338.vdb
%%RU%%drweb/bases/drw43339.vdb
%%RU%%drweb/bases/drw43340.vdb
%%RU%%drweb/bases/drw43341.vdb
drweb/bases/drwebase.vdb
drweb/bases/drwnasty.vdb
drweb/bases/drwrisky.vdb
drweb/bases/dwn43301.vdb
%%RU%%drweb/bases/dwn43302.vdb
%%RU%%drweb/bases/dwn43303.vdb
%%RU%%drweb/bases/dwr43301.vdb
drweb/bases/update.drl
@dirrm drweb/bases
@dirrm drweb/infected
@dirrm drweb/log
@dirrm drweb/run
@dirrm drweb/spool
@dirrm drweb/updates
@dirrm drweb
