# security/drweb/pkg-plist
@unexec if [ -f /var/drweb/run/drwebd.pid ] ; then %D/etc/rc.d/drwebd stop ; fi
drweb/agent/daemon.amc
drweb/drweb
drweb/drwebd
drweb/drwebdc
drweb/getkey.HOWTO
drweb/getkey.rus.HOWTO
drweb/getkey.rus.HOWTO.utf8
drweb/lib/drweb32.dll
drweb/lib/ru_daemon.dwl
drweb/lib/ru_scanner.dwl
drweb/read_signed
drweb/update/update.pl
@unexec if [ -f %D/etc/drweb/drweb32.ini ] && cmp -s %D/etc/drweb/drweb32.ini %D/etc/drweb/drweb32.ini-distr; then rm -f %D/etc/drweb/drweb32.ini ; fi
etc/drweb/drweb32.ini-distr
@exec if [ ! -f %B/drweb32.ini ] ; then cp -p %B/%F %B/drweb32.ini ; fi
@unexec if [ -f %D/etc/drweb/email.ini ] && cmp -s %D/etc/drweb/email.ini %D/etc/drweb/email.ini-distr; then rm -f %D/etc/drweb/email.ini ; fi
etc/drweb/email.ini-distr
@exec if [ ! -f %B/email.ini ] ; then cp -p %B/%F %B/email.ini ; fi
etc/drweb/monitor/daemon.mmc
etc/rc.d/drwebd
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog.rus
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/FAQ
%%PORTDOCS%%%%DOCSDIR%%/FAQ.rus
%%PORTDOCS%%%%DOCSDIR%%/FAQ.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/daemon/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/daemon/ChangeLog.rus
%%PORTDOCS%%%%DOCSDIR%%/daemon/ChangeLog.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/daemon/readme.daemon
%%PORTDOCS%%%%DOCSDIR%%/daemon/readme.daemon.rus
%%PORTDOCS%%%%DOCSDIR%%/daemon/readme.daemon.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/readme.eicar
%%PORTDOCS%%%%DOCSDIR%%/readme.eicar.rus
%%PORTDOCS%%%%DOCSDIR%%/readme.eicar.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/readme.license
%%PORTDOCS%%%%DOCSDIR%%/readme.license.rus
%%PORTDOCS%%%%DOCSDIR%%/readme.license.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/scanner/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/scanner/ChangeLog.rus
%%PORTDOCS%%%%DOCSDIR%%/scanner/ChangeLog.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/scanner/readme.scanner
%%PORTDOCS%%%%DOCSDIR%%/scanner/readme.scanner.rus
%%PORTDOCS%%%%DOCSDIR%%/scanner/readme.scanner.rus.utf8
%%PORTDOCS%%%%DOCSDIR%%/update/readme.update
%%PORTDOCS%%%%DOCSDIR%%/update/readme.update.rus
%%PORTDOCS%%%%DOCSDIR%%/update/readme.update.rus.utf8
%%PORTDOCS%%@dirrm %%DOCSDIR%%/update
%%PORTDOCS%%@dirrm %%DOCSDIR%%/scanner
%%PORTDOCS%%@dirrm %%DOCSDIR%%/daemon
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm etc/drweb/monitor
@dirrm etc/drweb
@dirrm drweb/update
@dirrm drweb/lib
@dirrm drweb/agent
@dirrm drweb
@cwd /var
drweb/bases/drw50000.vdb
drweb/bases/drw50001.vdb
drweb/bases/drw50002.vdb
drweb/bases/drw50003.vdb
drweb/bases/drw50004.vdb
drweb/bases/drw50005.vdb
drweb/bases/drw50006.vdb
drweb/bases/drw50007.vdb
drweb/bases/drw50008.vdb
drweb/bases/drw50009.vdb
drweb/bases/drw50010.vdb
drweb/bases/drw50011.vdb
drweb/bases/drw50012.vdb
drweb/bases/drw50013.vdb
drweb/bases/drw50014.vdb
drweb/bases/drw50015.vdb
drweb/bases/drw50016.vdb
drweb/bases/drw50017.vdb
drweb/bases/drw50018.vdb
drweb/bases/drw50019.vdb
drweb/bases/drw50020.vdb
drweb/bases/drw50021.vdb
drweb/bases/drw50022.vdb
drweb/bases/drw50023.vdb
drweb/bases/drw50024.vdb
drweb/bases/drw50025.vdb
drweb/bases/drw50026.vdb
drweb/bases/drw50027.vdb
drweb/bases/drw50028.vdb
drweb/bases/drw50029.vdb
drweb/bases/drw50030.vdb
drweb/bases/drw50031.vdb
drweb/bases/drwebase.vdb
drweb/bases/drwnasty.vdb
drweb/bases/drwrisky.vdb
drweb/bases/drwtoday.vdb
drweb/bases/dwn50001.vdb
drweb/bases/dwn50002.vdb
drweb/bases/dwntoday.vdb
drweb/bases/dwr50001.vdb
drweb/bases/dwrtoday.vdb
drweb/bases/update.drl
@dirrm drweb/bases
@dirrm drweb/cache
@dirrm drweb/infected
@dirrm drweb/log
@dirrm drweb/run
@dirrm drweb/spool
@dirrm drweb/updates
@dirrm drweb
