# security/drweb/pkg-plist
@comment $FreeBSD$
drweb/ChangeLog
drweb/README
drweb/README.RUS
drweb/bases/drw43001.vdb
drweb/bases/drw43002.vdb
drweb/bases/drw43003.vdb
drweb/bases/drw43004.vdb
drweb/bases/drw43005.vdb
drweb/bases/drw43006.vdb
drweb/bases/drw43007.vdb
drweb/bases/drw43008.vdb
drweb/bases/drw43009.vdb
drweb/bases/drwebase.vdb
drweb/doc/readme.daemon
drweb/doc/readme.daemon.rus
drweb/doc/readme.eicar
drweb/doc/readme.scanner
drweb/doc/readme.scanner.rus
drweb/drweb
@unexec if [ -f %D/drweb/drweb.key ] && cmp -s %D/drweb/drweb.key %D/drweb/drweb.key-dist; then rm -f %D/drweb/drweb.key; fi
drweb/drweb.key-dist
@exec if [ ! -f %D/drweb/drweb.key ] ; then cp -p %D/%F %B/drweb.key; fi
@unexec if [ -f %D/drweb/drweb32.ini ] && cmp -s %D/drweb/drweb32.ini %D/drweb/drweb32.ini-dist; then rm -f %D/drweb/drweb32.ini; fi
drweb/drweb32.ini-dist
@exec if [ ! -f %D/drweb/drweb32.ini ] ; then cp -p %D/%F %B/drweb32.ini; fi
drweb/drwebd
@unexec if [ -f %D/drweb/drwebd.key ] && cmp -s %D/drweb/drwebd.key %D/drweb/drwebd.key-dist; then rm -f %D/drweb/drwebd.key; fi
drweb/drwebd.key-dist
@exec if [ ! -f %D/drweb/drwebd.key ] ; then cp -p %D/%F %B/drwebd.key; fi
drweb/drwebdc
drweb/lib/drweb32.dll
drweb/lib/russian.dwl
drweb/update/readme.update
drweb/update/readme.update.rus
drweb/update/update.pl
etc/rc.d/00drwebd.sh
@dirrm drweb/updates
@dirrm drweb/update
@dirrm drweb/run
@dirrm drweb/lib
@dirrm drweb/infected
@dirrm drweb/doc
@dirrm drweb/bases
@dirrm drweb
