# sysutils/idled/files/idled.in
#!/bin/sh
#
# $FreeBSD: /tmp/pcvs/ports/sysutils/idled/files/Attic/idled.in,v 1.3 2010-03-27 00:14:56 dougb Exp $
#

# PROVIDE: idled
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable idled:
#
# idled_enable="YES"
# idled_flags="<set as needed>"
#
# See idled(8) for flags
#

. /etc/rc.subr

name=idled
rcvar=`set_rcvar`

command=%%PREFIX%%/libexec/idled
required_files=%%PREFIX%%/etc/idled.cf

# set defaults

idled_enable=${idled_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"
