# sysutils/msyslog/pkg-plist
etc/rc.d/msyslog.sh.sample
lib/alat/libmsyslog.so.1
sbin/peochk
sbin/syslogd
share/examples/msyslog/im_mymodule.c
share/examples/msyslog/om_mymodule.c
share/examples/msyslog/syslog.conf.classic
share/examples/msyslog/syslog.conf.mysql
share/examples/msyslog/syslog.conf.peo
share/examples/msyslog/syslog.conf.pgsql
share/examples/msyslog/syslog.conf.regex
@dirrm share/examples/msyslog
%%PORTDOCS%%share/doc/msyslog/AUTHORS
%%PORTDOCS%%share/doc/msyslog/COPYING
%%PORTDOCS%%share/doc/msyslog/HOW-TO-UPGRADE
%%PORTDOCS%%share/doc/msyslog/HOW_TO_WRITE_A_MODULE
%%PORTDOCS%%share/doc/msyslog/NEWS
%%PORTDOCS%%share/doc/msyslog/README
%%PORTDOCS%%share/doc/msyslog/TODO
%%PORTDOCS%%@dirrm share/doc/msyslog
@unexec rmdir %D/lib/alat 2>/dev/null || true
