# mail/courier-imap/files/openrc-courier-imap-imapd-ssl.in
#!/sbin/openrc-run

# Define these courier_imap_imapd_ssl_* variables in one of these files:
# 	/etc/rc.conf
# 	/etc/rc.conf.local
# 	/etc/rc.conf.d/courier_imap_imapd_ssl

name="courier_imap_imapd_ssl"
command="%%PREFIX%%/libexec/courier-imap/imapd-ssl.rc"
pidfile="/var/run/imapd-ssl.pid"
procname="%%PREFIX%%/sbin/courierlogger"

start() {
	case $1 in
	start)
		einfo "Starting ${name}."
		${command} start
		;;
	stop)
		einfo "Stopping ${name}."
		${command} stop
		;;
	esac
}

depend()
{
	keyword -shutdown
}
