# net/olsrd/files/olsrd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: olsrd
# REQUIRE: NETWORKING SERVERS

. /etc/rc.subr

name="olsrd"
rcvar=olsrd_enable

command=%%PREFIX%%/sbin/olsrd
command_args="-d 0"

olsrd_enable=${olsrd_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"

