# dns/nsd/files/nsd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: nsd
# REQUIRE: DAEMON
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable nsd:
#
# nsd_enable="YES"
#

. /etc/rc.subr

name=nsd
rcvar=${name}_enable

required_files=%%PREFIX%%/etc/nsd/nsd.conf

command=%%PREFIX%%/sbin/nsd-control
command_args="start"
pidfile=`%%PREFIX%%/sbin/nsd-checkconf -o pidfile %%PREFIX%%/etc/nsd/nsd.conf`
procname=%%PREFIX%%/sbin/${name}

load_rc_config ${name}

nsd_enable=${nsd_enable-"NO"}

extra_commands="reload"
reload_cmd="${name}_reload"
stop_cmd="${name}_stop"

nsd_reload()
{
	echo "Reloading ${name}."
	kill -HUP `cat $pidfile`
}

nsd_stop()
{
	echo "Stopping ${name}."
	kill -TERM `cat $pidfile`
}

run_rc_command "$1"

