# mail/pop3proxy/files/openrc-pop3proxyd.in
#!/sbin/openrc-run

name=pop3proxyd
command=%%PREFIX%%/sbin/pop3proxyd

depend() {
	keyword -shutdown
}
