# net/iaxmodem/files/iaxmodem.in
#!/bin/sh

# PROVIDE: iaxmodem
# REQUIRE: DAEMON
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable iaxmodem:
#
# iaxmodem_enable="YES"
#

. /etc/rc.subr

name=iaxmodem
rcvar=`set_rcvar`

command=%%PREFIX%%/sbin/iaxmodem

pidfile=${iaxmodem_pidfile:-"/var/run/iaxmodem.pid"}

iaxmodem_enable=${iaxmodem_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"
