# net/iaxmodem/files/openrc-iaxmodem.in
#!/sbin/openrc-run

name=iaxmodem
command=%%PREFIX%%/sbin/iaxmodem
pidfile=${iaxmodem_pidfile:-"/var/run/iaxmodem.pid"}

depend() {
	keyword -shutdown
}
