# www/elog/files/elogd.in
#!/bin/sh

# PROVIDE: elogd
# REQUIRE: network
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# elogd_enable (bool):   Set to NO by default.
#               Set it to YES to enable doormand.
# elogd_config (path):   Set to %%PREFIX%%/etc/elogd.cfg
#               by default.
# elogd_flags (args): See elogd(1) for possible options

. /etc/rc.subr

name="elogd"
rcvar=${name}_enable

command=%%PREFIX%%/sbin/${name}
pidfile=/var/run/${name}.pid

load_rc_config $name

: ${elogd_enable="NO"}
: ${elogd_config="%%PREFIX%%/etc/elogd.cfg"}

command_args="-D -c $elogd_config -s %%PREFIX%%/share/elog/ -f ${pidfile} -d %%PREFIX%%/%%LOGBOOKDIR%%"

run_rc_command "$1"
