# net/gnu-radius/files/radiusd.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE radiusd
# REQUIRE: NETWORKING SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable radiusd:
#
# radiusd_enable="YES"

. /etc/rc.subr

name=radiusd
rcvar=radiusd_enable

command=%%PREFIX%%/sbin/radiusd
radiusd_config=${radiusd_config:-"%%PREFIX%%/etc/raddb/config"}
required_files=${radiusd_config}
pidfile=/var/run/radiusd.pid

radiusd_enable=${radiusd_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"
