# net/bfilter/files/bfilter.sh.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: bfilter
# REQUIRE: DAEMON
#
# Add the following line to /etc/rc.conf to enable bfilter:
#
# bfilter_enable="YES"
#

. /etc/rc.subr

name=bfilter
rcvar=bfilter_enable

command=%%PREFIX%%/bin/${name}
sig_stop=KILL

load_rc_config ${name}

bfilter_enable=${bfilter_enable-"NO"}
bfilter_flags=${bfilter_flags-"-c %%PREFIX%%/etc/bfilter"}

run_rc_command "$1"
