# security/tang/pkg-plist
libexec/tangd
man/man1/tang-show-keys.1.gz
man/man8/tang.8.gz
libexec/tangd-keygen
libexec/tangd-rotate-keys
bin/tang-show-keys
etc/rc.d/tangd
@dir /var/db/tang
