# security/vault/files/openrc-vault.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# vault_user (user):        Set user to run vault.
#                               Default is "vault".
# vault_group (group):      Set group to run vault.
#                               Default is "vault".
# vault_config (dir):       Set vault config file.
#                               Default is "%%PREFIX%%/etc/vault.hcl".
# vault_env (dir):          Set environment variables used with vault

name=vault

: ${vault_user:="vault"}
: ${vault_group:="vault"}
: ${vault_config:="%%PREFIX%%/etc/vault.hcl"}
: ${vault_env:=""}

pidfile=/var/run/vault.pid
command="%%PREFIX%%/bin/vault"
command_args="server -config=${vault_config}"
supervisor=supervise-daemon

start_pre() {
        if [ ! -e ${pidfile} ]; then
                install -o ${vault_user} -g ${vault_group} /dev/null ${pidfile};
        fi

        if [ ! -d ${vault_dir} ]; then
                install -d -o ${vault_user} -g ${vault_group} ${vault_dir}
        fi
	/usr/bin/env ${vault_env} > /dev/null 
}

depend() {
	need localmount
	keyword -shutdown
}
