# mail/gld/pkg-plist
etc/gld.conf-dist
etc/rc.d/gld.sh
sbin/gld
%%PORTDOCS%%%%DOCSDIR%%/HISTORY
%%PORTDOCS%%%%DOCSDIR%%/LICENCE
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/table-whitelist.sql
%%PORTDOCS%%%%DOCSDIR%%/tables.mysql
%%PORTDOCS%%%%DOCSDIR%%/tables.pgsql
%%PORTDOCS%%@dirrm %%DOCSDIR%%
