# databases/mysql55-server/pkg-plist
@unexec %D/etc/rc.d/mysql-server.sh stop > /dev/null 2>&1 || true
bin/comp_err
bin/innochecksum
bin/my_print_defaults
bin/myisam_ftdump
bin/myisamchk
bin/myisamlog
bin/myisampack
bin/mysql_create_system_tables
bin/mysql_install_db
bin/mysql_tzinfo_to_sql
bin/mysql_upgrade
bin/mysql_waitpid
bin/mysqld_safe
%%NDB%%bin/ndb_config
%%NDB%%bin/ndb_delete_all
%%NDB%%bin/ndb_desc
%%NDB%%bin/ndb_drop_index
%%NDB%%bin/ndb_drop_table
%%NDB%%bin/ndb_error_reporter
%%NDB%%bin/ndb_mgm
%%NDB%%bin/ndb_print_backup_file
%%NDB%%bin/ndb_print_schema_file
%%NDB%%bin/ndb_print_sys_file
%%NDB%%bin/ndb_restore
%%NDB%%bin/ndb_select_all
%%NDB%%bin/ndb_select_count
%%NDB%%bin/ndb_show_tables
%%NDB%%bin/ndb_size.pl
%%NDB%%bin/ndb_test_platform
%%NDB%%bin/ndb_waiter
bin/perror
bin/replace
bin/resolve_stack_dump
bin/resolveip
%%NDB%%include/mysql/storage/ndb/mgmapi/mgmapi.h
%%NDB%%include/mysql/storage/ndb/mgmapi/mgmapi_config_parameters.h
%%NDB%%include/mysql/storage/ndb/mgmapi/mgmapi_config_parameters_debug.h
%%NDB%%include/mysql/storage/ndb/mgmapi/mgmapi_debug.h
%%NDB%%include/mysql/storage/ndb/mgmapi/ndb_logevent.h
%%NDB%%include/mysql/storage/ndb/mgmapi/ndbd_exit_codes.h
%%NDB%%include/mysql/storage/ndb/ndb_constants.h
%%NDB%%include/mysql/storage/ndb/ndb_init.h
%%NDB%%include/mysql/storage/ndb/ndb_types.h
%%NDB%%include/mysql/storage/ndb/ndb_version.h
%%NDB%%include/mysql/storage/ndb/ndbapi/Ndb.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbApi.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbBlob.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbDictionary.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbError.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbEventOperation.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbIndexOperation.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbIndexScanOperation.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbIndexStat.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbOperation.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbPool.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbRecAttr.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbReceiver.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbScanFilter.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbScanOperation.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/NdbTransaction.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/ndbapi_limits.h
%%NDB%%include/mysql/storage/ndb/ndbapi/ndb_cluster_connection.hpp
%%NDB%%include/mysql/storage/ndb/ndbapi/ndb_opt_defaults.h
%%NDB%%include/mysql/storage/ndb/ndbapi/ndberror.h
lib/mysql/libcsv.a
lib/mysql/libdbug.a
lib/mysql/libheap.a
lib/mysql/libmyisam.a
lib/mysql/libmyisammrg.a
lib/mysql/libmystrings.a
lib/mysql/libmysys.a
%%NDB%%lib/mysql/libndbclient.a
%%NDB%%lib/mysql/libndbclient.la
%%NDB%%lib/mysql/libndbclient.so
%%NDB%%lib/mysql/libndbclient.so.0
lib/mysql/libvio.a
libexec/mysqld
%%NDB%%libexec/ndb_cpcd
%%NDB%%libexec/ndb_mgmd
%%NDB%%libexec/ndbd
%%DATADIR%%/binary-configure
%%DATADIR%%/charsets/Index.xml
%%DATADIR%%/charsets/README
%%DATADIR%%/charsets/armscii8.xml
%%DATADIR%%/charsets/ascii.xml
%%DATADIR%%/charsets/cp1250.xml
%%DATADIR%%/charsets/cp1251.xml
%%DATADIR%%/charsets/cp1256.xml
%%DATADIR%%/charsets/cp1257.xml
%%DATADIR%%/charsets/cp850.xml
%%DATADIR%%/charsets/cp852.xml
%%DATADIR%%/charsets/cp866.xml
%%DATADIR%%/charsets/dec8.xml
%%DATADIR%%/charsets/geostd8.xml
%%DATADIR%%/charsets/greek.xml
%%DATADIR%%/charsets/hebrew.xml
%%DATADIR%%/charsets/hp8.xml
%%DATADIR%%/charsets/keybcs2.xml
%%DATADIR%%/charsets/koi8r.xml
%%DATADIR%%/charsets/koi8u.xml
%%DATADIR%%/charsets/latin1.xml
%%DATADIR%%/charsets/latin2.xml
%%DATADIR%%/charsets/latin5.xml
%%DATADIR%%/charsets/latin7.xml
%%DATADIR%%/charsets/macce.xml
%%DATADIR%%/charsets/macroman.xml
%%DATADIR%%/charsets/swe7.xml
%%DATADIR%%/czech/errmsg.sys
%%DATADIR%%/danish/errmsg.sys
%%DATADIR%%/dutch/errmsg.sys
%%DATADIR%%/english/errmsg.sys
%%DATADIR%%/errmsg.txt
%%DATADIR%%/estonian/errmsg.sys
%%DATADIR%%/fill_help_tables.sql
%%DATADIR%%/french/errmsg.sys
%%DATADIR%%/german/errmsg.sys
%%DATADIR%%/greek/errmsg.sys
%%DATADIR%%/hungarian/errmsg.sys
%%DATADIR%%/italian/errmsg.sys
%%DATADIR%%/japanese/errmsg.sys
%%DATADIR%%/korean/errmsg.sys
%%DATADIR%%/mi_test_all
%%DATADIR%%/mi_test_all.res
%%DATADIR%%/my-huge.cnf
%%DATADIR%%/my-innodb-heavy-4G.cnf
%%DATADIR%%/my-large.cnf
%%DATADIR%%/my-medium.cnf
%%DATADIR%%/my-small.cnf
%%DATADIR%%/mysql-log-rotate
%%DATADIR%%/mysql.server
%%DATADIR%%/ndb-config-2-node.ini
%%NDB%%%%DATADIR%%/ndb_size.tmpl
%%DATADIR%%/norwegian-ny/errmsg.sys
%%DATADIR%%/norwegian/errmsg.sys
%%DATADIR%%/polish/errmsg.sys
%%DATADIR%%/portuguese/errmsg.sys
%%DATADIR%%/romanian/errmsg.sys
%%DATADIR%%/russian/errmsg.sys
%%DATADIR%%/serbian/errmsg.sys
%%DATADIR%%/slovak/errmsg.sys
%%DATADIR%%/spanish/errmsg.sys
%%DATADIR%%/swedish/errmsg.sys
%%DATADIR%%/ukrainian/errmsg.sys
%%NDB%%@dirrm include/mysql/storage/ndb/mgmapi
%%NDB%%@dirrm include/mysql/storage/ndb/ndbapi
%%NDB%%@dirrm include/mysql/storage/ndb
%%NDB%%@dirrm include/mysql/storage
@dirrmtry include/mysql
@dirrmtry lib/mysql
@dirrm %%DATADIR%%/charsets
@dirrm %%DATADIR%%/czech
@dirrm %%DATADIR%%/danish
@dirrm %%DATADIR%%/dutch
@dirrm %%DATADIR%%/english
@dirrm %%DATADIR%%/estonian
@dirrm %%DATADIR%%/french
@dirrm %%DATADIR%%/german
@dirrm %%DATADIR%%/greek
@dirrm %%DATADIR%%/hungarian
@dirrm %%DATADIR%%/italian
@dirrm %%DATADIR%%/japanese
@dirrm %%DATADIR%%/korean
@dirrm %%DATADIR%%/norwegian
@dirrm %%DATADIR%%/norwegian-ny
@dirrm %%DATADIR%%/polish
@dirrm %%DATADIR%%/portuguese
@dirrm %%DATADIR%%/romanian
@dirrm %%DATADIR%%/russian
@dirrm %%DATADIR%%/serbian
@dirrm %%DATADIR%%/slovak
@dirrm %%DATADIR%%/spanish
@dirrm %%DATADIR%%/swedish
@dirrm %%DATADIR%%/ukrainian
@dirrmtry %%DATADIR%%
