# sysutils/dsbmd/files/openrc-dsbmd.in
#!/sbin/openrc-run

name=dsbmd

command=%%PREFIX%%/libexec/${name}
pidfile=/var/run/${name}.pid

depend() {
	use localmount
	keyword -shutdown
}
