# multimedia/ffmpeg23/files/ffserver23.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ffserver23
# REQUIRE: NETWORKING
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable ffserver:
#
#ffserver23_enable="YES"
#
ffserver23_enable="${ffserver_enable-NO}"

. /etc/rc.subr

name=ffserver23
rcvar=ffserver23_enable

command="%%PREFIX%%/bin/${name}"
command_args="&"
required_files=%%PREFIX%%/etc/ffserver23.conf

load_rc_config ${name}
run_rc_command "$1"
