# security/openvpn/files/openrc-openvpn.in
#!/sbin/openrc-run

name=openvpn
command="%%PREFIX%%/sbin/openvpn"
pidfile="/var/run/openvpn.pid"
vpnconf="%%PREFIX%%/etc/openvpn/openvpn.conf"
required_files=$vpnconf
supervisor=supervise-daemon
command_args="--config ${vpnconf} --writepid ${pidfile} --daemon ${name}"

depend() {
	need localmount net
	before netmount
	after bootmisc
}

start_pre() {
	# If the config file does not specify the cd option, we do
	# But if we specify it, we override the config option which we do not want
	if ! grep -q "^[ 	]*cd[ 	].*" "${vpnconf}" ; then
		#$command_args+="${args} --cd %%PREFIX%%/etc/openvpn"
	fi
}
