# print/cloud-print-connector/files/gcp-cups-connectord.in
#!/bin/sh
#
# PROVIDE: gcp_cups_connectord
# REQUIRE: DAEMON NETWORKING
# KEYWORD: shutdown
#
# gcp_cups_connectord_enable (bool):  Set to "NO" by default.
#                                     Set it to "YES" to enable
#

. /etc/rc.subr

name="gcp_cups_connectord"
rcvar="gcp_cups_connectord_enable"

load_rc_config $name

: ${name_enable:=no}

pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
bin="%%PREFIX%%/bin/gcp-cups-connector"
start_precmd="${name}_prestart"

gcp_cups_connectord_prestart() {
  if [ -z "${rc_flags}" ];then
    rc_flags="-P ${pidfile} -f ${bin} -config-filename %%PREFIX%%/etc/gcp-cups-connector.config.json"
  else
    rc_flags="-P ${pidfile} -f ${bin} ${rc_flags}"
  fi
  return 0
}

run_rc_command "$1"
