# net/dgd/pkg-plist
bin/precomp
etc/dgd/kernel.dgd
etc/rc.d/dgd.sh.sample
sbin/dgd
%%DATADIR%%/kernel/doc/kernel/efun/add_event
%%DATADIR%%/kernel/doc/kernel/efun/call_limited
%%DATADIR%%/kernel/doc/kernel/efun/call_other
%%DATADIR%%/kernel/doc/kernel/efun/call_trace
%%DATADIR%%/kernel/doc/kernel/efun/clone_object
%%DATADIR%%/kernel/doc/kernel/efun/compile_object
%%DATADIR%%/kernel/doc/kernel/efun/destruct_object
%%DATADIR%%/kernel/doc/kernel/efun/event
%%DATADIR%%/kernel/doc/kernel/efun/event_except
%%DATADIR%%/kernel/doc/kernel/efun/file_info
%%DATADIR%%/kernel/doc/kernel/efun/find_object
%%DATADIR%%/kernel/doc/kernel/efun/get_dir
%%DATADIR%%/kernel/doc/kernel/efun/new_object
%%DATADIR%%/kernel/doc/kernel/efun/query_events
%%DATADIR%%/kernel/doc/kernel/efun/query_subscribed_event
%%DATADIR%%/kernel/doc/kernel/efun/remove_event
%%DATADIR%%/kernel/doc/kernel/efun/status
%%DATADIR%%/kernel/doc/kernel/efun/subscribe_event
%%DATADIR%%/kernel/doc/kernel/efun/unsubscribe_event
%%DATADIR%%/kernel/doc/kernel/hook/driver
%%DATADIR%%/kernel/doc/kernel/hook/rsrcd
%%DATADIR%%/kernel/doc/kernel/hook/userd
%%DATADIR%%/kernel/doc/kernel/lfun/allow_subscribe
%%DATADIR%%/kernel/doc/kernel/lfun/create
%%DATADIR%%/kernel/doc/kernel/lfun/query_owner
%%DATADIR%%/kernel/doc/kernel/overview
%%DATADIR%%/kernel/include/config.h
%%DATADIR%%/kernel/include/kernel/access.h
%%DATADIR%%/kernel/include/kernel/kernel.h
%%NETWORKING%%%%DATADIR%%/kernel/include/kernel/net.h
%%DATADIR%%/kernel/include/kernel/objreg.h
%%DATADIR%%/kernel/include/kernel/rsrc.h
%%DATADIR%%/kernel/include/kernel/tls.h
%%DATADIR%%/kernel/include/kernel/user.h
%%DATADIR%%/kernel/include/kernel/version.h
%%DATADIR%%/kernel/include/std.h
%%DATADIR%%/kernel/kernel/lib/api/access.c
%%DATADIR%%/kernel/kernel/lib/api/objreg.c
%%DATADIR%%/kernel/kernel/lib/api/rsrc.c
%%DATADIR%%/kernel/kernel/lib/api/tls.c
%%DATADIR%%/kernel/kernel/lib/api/user.c
%%DATADIR%%/kernel/kernel/lib/auto.c
%%DATADIR%%/kernel/kernel/lib/connection.c
%%NETWORKING%%%%DATADIR%%/kernel/kernel/lib/network/port.c
%%DATADIR%%/kernel/kernel/lib/user.c
%%DATADIR%%/kernel/kernel/lib/wiztool.c
%%DATADIR%%/kernel/kernel/obj/binary.c
%%NETWORKING%%%%DATADIR%%/kernel/kernel/obj/port.c
%%DATADIR%%/kernel/kernel/obj/rsrc.c
%%DATADIR%%/kernel/kernel/obj/telnet.c
%%DATADIR%%/kernel/kernel/obj/user.c
%%DATADIR%%/kernel/kernel/obj/wiztool.c
%%DATADIR%%/kernel/kernel/sys/accessd.c
%%NETWORKING%%%%DATADIR%%/kernel/kernel/sys/binary_port.c
%%DATADIR%%/kernel/kernel/sys/driver.c
%%NETWORKING%%%%DATADIR%%/kernel/kernel/sys/emergency_port.c
%%DATADIR%%/kernel/kernel/sys/objregd.c
%%DATADIR%%/kernel/kernel/sys/rsrcd.c
%%NETWORKING%%%%DATADIR%%/kernel/kernel/sys/telnet_port.c
%%DATADIR%%/kernel/kernel/sys/userd.c
%%NETWORKING%%%%DATADIR%%/kernel/usr/System/initd.c
@unexec rm -f %D/%%DATADIR%%/kernel/kernel/data/*
@unexec rm -f %D/%%DATADIR%%/kernel/include/float.h
@unexec rm -f %D/%%DATADIR%%/kernel/include/limits.h
@unexec rm -f %D/%%DATADIR%%/kernel/include/status.h
@unexec rm -f %D/%%DATADIR%%/kernel/include/trace.h
@unexec rm -f %D/%%DATADIR%%/kernel/include/type.h
@unexec rm -f %D/%%DATADIR%%/kernel.log
@unexec rm -f %D/%%DATADIR%%/tmp/dgd.swap
@dirrm %%DATADIR%%/tmp
@dirrm %%DATADIR%%/kernel/usr/admin
@dirrm %%DATADIR%%/kernel/usr/System
@dirrm %%DATADIR%%/kernel/usr
@dirrm %%DATADIR%%/kernel/kernel/sys
@dirrm %%DATADIR%%/kernel/kernel/obj
%%NETWORKING%%@dirrm %%DATADIR%%/kernel/kernel/lib/network
@dirrm %%DATADIR%%/kernel/kernel/lib/api
@dirrm %%DATADIR%%/kernel/kernel/lib
@dirrm %%DATADIR%%/kernel/kernel/data
@dirrm %%DATADIR%%/kernel/kernel
@dirrm %%DATADIR%%/kernel/include/kernel
@dirrm %%DATADIR%%/kernel/include
@dirrm %%DATADIR%%/kernel/doc/kernel/lfun
@dirrm %%DATADIR%%/kernel/doc/kernel/hook
@dirrm %%DATADIR%%/kernel/doc/kernel/efun
@dirrm %%DATADIR%%/kernel/doc/kernel
@dirrm %%DATADIR%%/kernel/doc
@dirrm %%DATADIR%%/kernel
@dirrm %%DATADIR%%
@dirrm etc/dgd
