# net/dgd/pkg-plist
bin/precomp
etc/dgd/mud.dgd.sample
etc/rc.d/dgd.sh.sample
sbin/dgd
@dirrm etc/dgd
