# dns/rbldnsd/pkg-plist
sbin/rbldnsd
etc/rbldnsd/example
etc/rc.d/rbldnsd.sh
%%PORTDOCS%%%%DOCSDIR%%/CHANGES-0.81
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%%%DOCSDIR%%/changelog
%%PORTDOCS%%%%DOCSDIR%%/copyright
%%PORTDOCS%%%%DOCSDIR%%/rbldnsd.default
%%PORTDOCS%%%%DOCSDIR%%/rbldnsd.init
%%PORTDOCS%%%%DOCSDIR%%/rbldnsd.spec
%%PORTDOCS%%%%DOCSDIR%%/rbldnsd.html
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@unexec rmdir %D/etc/rbldnsd/ 2>/dev/null || true
