# www/drood/pkg-plist
bin/drood
bin/scgid
etc/rc.d/drood
etc/rc.d/scgid
@unexec if cmp -s %D/etc/drood/types.tab %D/etc/drood/types.tab.dist; then rm -f %D/etc/drood/types.tab; fi
etc/drood/types.tab.dist
@dirrm /etc/drood
