# net/flowscan/pkg-plist
etc/rc.d/cflowd-flowscan.sh.sample
var/db/flows/bin/flowscan
var/db/flows/bin/FlowScan.pm
var/db/flows/bin/CampusIO.pm
var/db/flows/bin/SubNetIO.pm
var/db/flows/bin/locker
var/db/flows/bin/add_ds.pl
var/db/flows/bin/add_txrx
var/db/flows/bin/event2vrule
var/db/flows/bin/ip2hostname
var/db/flows/bin/CampusIO.cf.sample
var/db/flows/bin/Napster_subnets.boulder.sample
var/db/flows/bin/SubNetIO.cf.sample
var/db/flows/bin/flowscan.cf.sample
var/db/flows/bin/local_nets.boulder.sample
var/db/flows/bin/our_subnets.boulder.sample
var/db/flows/graphs/Makefile.sample
share/doc/flowscan/README
share/doc/flowscan/README.html
share/doc/flowscan/INSTALL
share/doc/flowscan/INSTALL.html
share/doc/flowscan/TODO
share/doc/flowscan/crontab.sample
share/doc/flowscan/cflowd.sample
share/doc/flowscan/flowscan.sample
@dirrm var/db/flows/graphs
@dirrm var/db/flows/bin
@dirrm var/db/flows
@unexec rmdir %D/var/db 2>/dev/null || true
@unexec rmdir %D/var 2>/dev/null || true
@dirrm share/doc/flowscan
