# net/dante/files/sockd.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: sockd
# REQUIRE: NETWORKING
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable sockd:
#
#sockd_enable="YES"
#

. /etc/rc.subr

name=sockd
rcvar=sockd_enable

command=%%PREFIX%%/sbin/sockd
pidfile=/var/run/sockd.pid
required_files=%%PREFIX%%/etc/sockd.conf

sockd_enable=${sockd_enable:-"NO"}
sockd_flags=${sockd_flags:-"-D"}

load_rc_config $name
run_rc_command "$1"
