# net/ssltunnel-server/pkg-plist
etc/rc.d/ssltunnel.sh.sample
@unexec if cmp -s %D/etc/ssltunnel/tunnel.conf %D/etc/ssltunnel/tunnel.conf.default; then rm -f %D/etc/ssltunnel/tunnel.conf; fi
@unexec if cmp -s %D/etc/ssltunnel/users %D/etc/ssltunnel/users.sample; then rm -f %D/etc/ssltunnel/users; fi
etc/ssltunnel/tunnel.conf.default
etc/ssltunnel/users.sample
libexec/pppserver
sbin/pppwho
%%PORTDOCS%%%%DOCSDIR%%/LISEZ-MOI
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm etc/ssltunnel
