# sysutils/tenshi/files/openrc-tenshi.in
#!/sbin/openrc-run

name='tenshi'

command='%%PREFIX%%/bin/tenshi'
pidfile='/var/run/tenshi.pid'
required_files='%%PREFIX%%/etc/tenshi.conf'

depend() {
	keyword -shutdown
}
