# net/pfflowd/files/pfflowd.in
#!/bin/sh
# $FreeBSD$
#

# PROVIDE: pfflowd
# REQUIRE: LOGIN

#
# Add the following lines to /etc/rc.conf to enable pfflowd:
#
# pfflowd_enable="YES"
# pfflowd_host="<host>:<port>"
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#

. /etc/rc.subr

name=pfflowd
rcvar=`set_rcvar`

load_rc_config $name

# set defaults

pfflowd_enable=${pfflowd_enable:-"NO"}
pfflowd_host=${pfflowd_host:-"127.0.0.1:2055"}

command=%%PREFIX%%/sbin/pfflowd
command_args="-n ${pfflowd_host}"
extra_commands=reload

stop_postcmd=stop_postcmd

stop_postcmd()
{
  rm -f $pidfile
}

pidfile="/var/run/${name}.pid"

run_rc_command "$1"
