# net/zebra-devel/pkg-plist
sbin/bgpd
sbin/ospfd
sbin/ripd
sbin/zebra
etc/rc.d/zebra.sh
etc/zebra/bgpd.conf.sample
etc/zebra/ripd.conf.sample
etc/zebra/zebra.conf.sample
etc/zebra/ospfd.conf.sample
@unexec install-info --delete %D/info/zebra.info %D/info/dir
info/zebra.info
@exec install-info %D/info/zebra.info %D/info/dir
@dirrm etc/zebra
