# net/poptop/pkg-plist
etc/pptpd.conf.sample
etc/rc.d/pptpd.sh.sample
lib/pptpd/pptpd-logwtmp.so
sbin/pptpctrl
sbin/pptpd
@dirrm lib/pptpd
