# net/poptop/pkg-plist
@unexec %D/etc/rc.d/pptpd.sh stop > /dev/null 2>&1 || true
sbin/pptpd
sbin/pptpctrl
lib/pptpd/pptpd-logwtmp.so
etc/pptpd.conf.sample
@dirrmtry lib/pptpd
