# mail/pop-before-smtp/pkg-plist
@comment $FreeBSD$
etc/pop-before-smtp-conf.pl.sample
etc/rc.d/pop-before-smtp.sh.sample
sbin/pop-before-smtp
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README.QUICKSTART
%%PORTDOCS%%%%DOCSDIR%%/README.rootless-install
%%PORTDOCS%%%%DOCSDIR%%/popauth.m4
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%%%DOCSDIR%%/popa3d/README
%%PORTDOCS%%%%DOCSDIR%%/popa3d/popa3d-0.4.patch
%%PORTDOCS%%%%DOCSDIR%%/popa3d/popa3d-0.6.4.patch
%%PORTDOCS%%@dirrm %%DOCSDIR%%/popa3d
%%PORTDOCS%%@dirrm %%DOCSDIR%%
