# sysutils/pcbsd-syscache/pkg-plist
bin/syscache
bin/syscache-daemon
etc/rc.d/syscache
