# security/ssh2/pkg-plist
bin/ssh2
bin/scp2
bin/sftp2
bin/ssh-agent2
%%WITH_X11:%%bin/ssh-askpass2
bin/ssh-keygen2
bin/ssh-add2
bin/ssh-signer2
bin/ssh-probe2
bin/sftp-server2
bin/ssh-dummy-shell
bin/ssh-pubkeymgr
bin/ssh-chrootmgr
bin/ssh
bin/ssh-agent
bin/ssh-add
bin/ssh-askpass
bin/ssh-keygen
bin/scp
bin/sftp
bin/sftp-server
bin/ssh-signer
bin/ssh-probe
etc/rc.d/sshd.sh
etc/ssh2/sshd2_config.sample
etc/ssh2/ssh2_config.sample
etc/ssh2/ssh_dummy_shell.out
sbin/sshd2
sbin/sshd-check-conf
sbin/sshd
@exec if [ ! -f %D/etc/ssh2/hostkey ]; then umask 022; echo "Generating 1024 bit host key."; %D/bin/ssh-keygen2 -P -b 1024 -t dsa %D/etc/ssh2/hostkey;  fi
%%PORTDOCS%%share/doc/ssh2/CHANGES
%%PORTDOCS%%share/doc/ssh2/FAQ
%%PORTDOCS%%share/doc/ssh2/INSTALL
%%PORTDOCS%%share/doc/ssh2/LICENSE
%%PORTDOCS%%share/doc/ssh2/MANIFEST
%%PORTDOCS%%share/doc/ssh2/NEWS
%%PORTDOCS%%share/doc/ssh2/README
%%PORTDOCS%%share/doc/ssh2/REGEX-SYNTAX
%%PORTDOCS%%share/doc/ssh2/SSH2.QUICKSTART
%%PORTDOCS%%@dirrm share/doc/ssh2
