# security/ssh2/pkg-plist
bin/ssh2
bin/scp2
bin/sftp2
bin/ssh-agent2
%%WITH_X11:%%bin/ssh-askpass2
bin/ssh-keygen2
bin/ssh-add2
bin/ssh-signer2
bin/ssh-probe2
bin/sftp-server2
%%STATIC%%bin/sftp-server2.static
bin/ssh-dummy-shell
%%STATIC%%bin/ssh-dummy-shell.static
bin/ssh
bin/ssh-agent
bin/ssh-add
%%WITH_X11:%%bin/ssh-askpass
bin/ssh-keygen
bin/scp
bin/sftp
bin/sftp-server
bin/ssh-signer
bin/ssh-probe
etc/rc.d/sshd.sh
@unexec if cmp -s %D/etc/ssh2/sshd2_config %D/etc/ssh2/sshd2_config.example; then rm -f %D/etc/ssh2/sshd2_config; fi
etc/ssh2/sshd2_config.example
@exec [ -f %B/sshd2_config ] || cp %B/%f %B/sshd2_config
@unexec if cmp -s %D/etc/ssh2/ssh2_config %D/etc/ssh2/ssh2_config.example; then rm -f %D/etc/ssh2/ssh2_config; fi
etc/ssh2/ssh2_config.example
@exec [ -f %B/ssh2_config ] || cp %B/%f %B/ssh2_config
etc/ssh2/ssh_dummy_shell.out
etc/ssh2/subconfig/anonymous.example
etc/ssh2/subconfig/host_ext.example
etc/ssh2/subconfig/host_int.example
etc/ssh2/subconfig/user.example
@exec [ -d %D/etc/ssh2/hostkeys ] || mkdir %D/etc/ssh2/hostkeys
@exec [ -d %D/etc/ssh2/knownhosts ] || mkdir %D/etc/ssh2/knownhosts
sbin/sshd2
sbin/sshd-check-conf
sbin/sshd
@exec if [ ! -f %D/etc/ssh2/hostkey ]; then umask 022; echo "Generating host key."; %D/bin/ssh-keygen2 -P -t dsa "DSA hostkey" %D/etc/ssh2/hostkey;  fi
%%PORTDOCS%%%%DOCSDIR%%/CHANGES
%%PORTDOCS%%%%DOCSDIR%%/FAQ
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/LICENSE
%%PORTDOCS%%%%DOCSDIR%%/MANIFEST
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/REGEX-SYNTAX
%%PORTDOCS%%%%DOCSDIR%%/SSH2.QUICKSTART
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@unexec rmdir %D/etc/ssh2/hostkeys 2> /dev/null || true
@unexec rmdir %D/etc/ssh2/knownhosts 2> /dev/null || true
@unexec rmdir %D/etc/ssh2/subconfig 2> /dev/null || true
@unexec rmdir %D/etc/ssh2 2> /dev/null || echo "If permanently deleting this package, %D/etc/ssh2 and its contents must be removed manually."
