# net/gnu-dico/files/openrc-dicod.in
#!/sbin/openrc-run

name=dicod
command=%%PREFIX%%/sbin/dicod
