# mail/pop3vscan/pkg-plist
etc/pop3vscan.conf.sample
etc/pop3vscan.mail.sample
etc/rc.d/pop3vscan.sh
sbin/pop3vscan
