# net/ocserv/files/openrc-ocserv.in
#!/sbin/openrc-run

name="ocserv"

: ${ocserv_pidfile:=/var/run/${name}/pid}

command=/usr/local/sbin/${name}
pidfile=${ocserv_pidfile}

depend() {
	keyword -shutdown
}
