# sysutils/cmogstored/files/openrc-cmogstored.in
#!/sbin/openrc-run

name=cmogstored

: ${cmogstored_config:="%%PREFIX%%/etc/mogstored.conf"}

required_files="${cmogstored_config}"

command="%%PREFIX%%/sbin/${name}"
command_args="--config='${cmogstored_config}' -d"

depend() {
	keyword -shutdown
}
