# sysutils/p5-MogileFS-Server/files/mogilefsd.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: mogilefsd
# REQUIRE: LOGIN cleanvar
# KEYWORD: shutdown

. /etc/rc.subr

name="mogilefsd"
rcvar=mogilefsd_enable
required_files=%%PREFIX%%/etc/mogtracker.conf

[ -z "$mogilefsd_enable" ] && mogilefsd_enable="NO"
[ -z "$mogilefsd_flags" ] && mogilefsd_flags="--config=%%PREFIX%%/etc/mogtracker.conf"

run_rc_command "$1"
