# sysutils/p5-MogileFS-Server/files/openrc-mogilefsd.in
#!/sbin/openrc-run

name=mogilefsd

required_files=%%PREFIX%%/etc/mogtracker.conf

command="%%PREFIX%%/bin/${name}"
command_user=nobody
command_args="--config=%%PREFIX%%/etc/mogtracker.conf"

depend() {
	keyword -shutdown
}
