# sysutils/p5-MogileFS-Server/files/openrc-mogstored.in
#!/sbin/openrc-run

name=mogstored
required_files=%%PREFIX%%/etc/mogstored.conf

command="%%PREFIX%%/bin/${name}"
command_args="--config=%%PREFIX%%/etc/mogstored.conf -d"

depend() {
	keyword -shutdown
}
