# japanese/skk/pkg-plist
bin/skkdic-expr
bin/skkdic-count
bin/skkdic-sort
libexec/skkserv
@unexec install-info --delete --section="The Emacs editor and associated tools" --entry="* Skk: (skk.info).		Yet another Japanese Inputting Method. (Japanese )" %D/info/skk.info %D/info/dir
info/skk.info
info/skk.info-1
info/skk.info-2
info/skk.info-3
info/skk.info-4
info/skk.info-5
info/skk.info-6
@exec [ -f %D/info/dir ] || sed -ne '1,/Menu:/p' /usr/share/info/dir > %D/info/dir
@exec install-info --section="The Emacs editor and associated tools" --entry="* Skk: (skk.info).		Yet another Japanese Inputting Method. (Japanese )" %D/info/skk.info %D/info/dir
share/emacs/site-lisp/avltree.el
share/emacs/site-lisp/avltree.elc
share/emacs/site-lisp/elib-node.el
share/emacs/site-lisp/elib-node.elc
share/emacs/site-lisp/queue-m.el
share/emacs/site-lisp/queue-m.elc
share/emacs/site-lisp/skk-auto.el
share/emacs/site-lisp/skk-auto.elc
share/emacs/site-lisp/skk-comp.el
share/emacs/site-lisp/skk-comp.elc
share/emacs/site-lisp/skk-e19.el
share/emacs/site-lisp/skk-e19.elc
share/emacs/site-lisp/skk-gadget.el
share/emacs/site-lisp/skk-gadget.elc
share/emacs/site-lisp/skk-isearch.el
share/emacs/site-lisp/skk-isearch.elc
share/emacs/site-lisp/skk-kakasi.el
share/emacs/site-lisp/skk-kakasi.elc
share/emacs/site-lisp/skk-kcode.el
share/emacs/site-lisp/skk-kcode.elc
share/emacs/site-lisp/skk-num.el
share/emacs/site-lisp/skk-num.elc
share/emacs/site-lisp/skk-server.el
share/emacs/site-lisp/skk-server.elc
share/emacs/site-lisp/skk-tree.el
share/emacs/site-lisp/skk-tree.elc
share/emacs/site-lisp/skk-tut.el
share/emacs/site-lisp/skk-tut.elc
share/emacs/site-lisp/skk-vip.el
share/emacs/site-lisp/skk-vip.elc
share/emacs/site-lisp/skk.el
share/emacs/site-lisp/skk.elc
share/emacs/site-lisp/stack-m.el
share/emacs/site-lisp/stack-m.elc
share/emacs/site-lisp/string.el
share/emacs/site-lisp/string.elc
share/skk/SKK.tut
share/skk/SKK.tut.E
share/skk/SKK-JISYO.L
etc/rc.d/skkserv.sh
@dirrm share/skk
