# japanese/skk/pkg-plist
bin/SKK/skkserv
bin/SKK/skkjisyo-entries
bin/SKK/skkjisyo-expr
bin/SKK/skkjisyo-sort
etc/rc.d/skkserv.sh
info/skk
lib/mule/site-lisp/skk-isearch.el
lib/mule/site-lisp/skk-isearch.elc
lib/mule/site-lisp/skk-tools.el
lib/mule/site-lisp/skk-tools.elc
lib/mule/site-lisp/skk-tut.el
lib/mule/site-lisp/skk-tut.elc
lib/mule/site-lisp/skk.el
lib/mule/site-lisp/skk.elc
lib/mule/19.28/i386--freebsd/emacs-sort
lib/mule/19.28/etc/SKK.tut
lib/mule/19.28/etc/SKK.tut.E
lib/skk/SKK-JISYO.JIS2
lib/skk/SKK-JISYO.L
lib/skk/SKK-JISYO.M
lib/skk/SKK-JISYO.S
lib/skk/latexinfo.sty
lib/skk/skk-romkan.tex
lib/skk/skk.ps
lib/skk/skk.tex
lib/skk/tabular.sty
lib/skk/dict/SKK-JISYO.L
