# devel/distcc/files/distccd.in
#!/bin/sh
#
# $FreeBSD$
#    $MCom: ports/devel/distcc/files/distccd.sh.in,v 1.1 2007/05/19 17:43:14 ahze Exp $

# PROVIDE: distccd
# REQUIRE: NETWORKING ldconfig
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable distccd:
#
# distccd_enable="YES"
#

distccd_enable=${distccd_enable-"NO"}
distccd_flags=${distccd_flags-"-a 127.0.0.0/8 --user distcc --daemon -P %%DISTCCD_PIDFILE%%"}

. /etc/rc.subr

name=distccd
rcvar=distccd_enable

command=%%PREFIX%%/sbin/${name}
# extract pid file from distccd_flags (if any)
pidfile="`echo ${distccd_flags} |sed -nE -e 's/^.*-P ([^ ]+).*$/\1/p'`"
# set default value
pidfile=${pidfile:-%%DISTCCD_PIDFILE%%}

distccd_precmd()
{
	# distccd drops root privileges and then attempts to create the pid file
	touch ${pidfile}
	chown distcc:distcc ${pidfile}
}

start_precmd="distccd_precmd"

load_rc_config ${name}
run_rc_command "$1"
