# net/isc-dhcp30-server/pkg-plist
@comment $FreeBSD$
etc/dhclient.conf.sample
etc/dhcpd.conf.sample
etc/rc.d/isc-dhcpd.sh.sample
include/dhcpctl.h
include/isc/boolean.h
include/isc/dst.h
include/isc/int.h
include/isc/lang.h
include/isc/list.h
include/isc/result.h
include/isc/types.h
include/omapip/alloc.h
include/omapip/buffer.h
include/omapip/omapip.h
lib/libdhcpctl.a
lib/libomapi.a
sbin/dhclient
sbin/dhclient-script
sbin/dhcpd
sbin/dhcrelay
share/doc/isc-dhcp3/CHANGES
share/doc/isc-dhcp3/COPYRIGHT
share/doc/isc-dhcp3/README
share/doc/isc-dhcp3/RELNOTES
@dirrm include/isc
@dirrm include/omapip
@dirrm share/doc/isc-dhcp3
@exec [ -f /var/db/dhclient.leases ] || touch /var/db/dhclient.leases
@unexec [ -s /var/db/dhclient.leases ] || rm -f /var/db/dhclient.leases
@exec [ -f /var/db/dhcpd.leases ] || touch /var/db/dhcpd.leases
@unexec [ -s /var/db/dhcpd.leases ] || rm -f /var/db/dhcpd.leases
@exec [ -f %D/etc/dhclient.conf ] || touch %D/etc/dhclient.conf
@unexec [ -L %D/etc/dhclient.conf -o -s %D/etc/dhclient.conf ] || rm -f %D/etc/dhclient.conf
