# graphics/airsaned/pkg-plist
sbin/airsaned
etc/airsane/Gnome-scanner.png
etc/rc.d/airsaned
@sample etc/airsane/ignore.conf.sample
@sample etc/airsane/options.conf.sample
