# security/softether-devel/files/openrc-vpnbridge.in
#!/sbin/openrc-run

name=vpnbridge
command=%%PREFIX%%/sbin/${name}
command_args="start"

stop_pre() {
	$command stop
}
