# net-mgmt/netustad/files/netustad.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: netustad
# REQUIRE: NETWORKING SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# netustad_enable (bool):	Set to NO by default.
#				Set it to YES to enable netustad.

. /etc/rc.subr

name="netustad"
rcvar=${name}_enable

command=%%PREFIX%%/sbin/${name}

load_rc_config $name

: ${netustad="NO"}

run_rc_command "$1"

