# net-mgmt/pftabled/files/openrc-pftabled.in
#!/sbin/openrc-run
#
# Add the following lines to /etc/rc.conf to enable pftabled:
# pftabled_flags (flags):	Set flags to alter default behaviour of pftabled.
# 				Consult with the man page for more info.

name="pftabled"

: ${pftabled_flags="-d"}
command="%%PREFIX%%/sbin/${name}"
command_args=${pftabled_flags}

depend() {
	use net
	keyword -shutdown
}
