# sysutils/syslog-ng3/pkg-plist
@unexec %D/etc/rc.d/syslog-ng.sh stop > /dev/null 2>&1 || true
etc/modules.conf
etc/scl.conf
@exec test -f %D/etc/syslog-ng.conf || cp -p %D/etc/syslog-ng.conf.sample %D/etc/syslog-ng.conf
@unexec if cmp -s %D/etc/syslog-ng.conf %D/etc/syslog-ng.conf.sample; then rm -f %D/etc/syslog-ng.conf; fi
etc/syslog-ng.conf.dist
etc/syslog-ng.conf.sample
bin/update-patterndb
bin/loggen
bin/pdbtool
lib/libsyslog-ng.la
lib/libsyslog-ng.so
lib/libsyslog-ng.so.0
lib/libsyslog-ng-patterndb.la
lib/libsyslog-ng-patterndb.so
lib/libsyslog-ng-patterndb.so.0
lib/syslog-ng/libaffile.la
lib/syslog-ng/libaffile.so
lib/syslog-ng/libafprog.la
lib/syslog-ng/libafprog.so
lib/syslog-ng/libafsocket.so
lib/syslog-ng/libafsocket-notls.la
lib/syslog-ng/libafsocket-notls.so
lib/syslog-ng/libafuser.la
lib/syslog-ng/libafuser.so
lib/syslog-ng/libbasicfuncs.la
lib/syslog-ng/libbasicfuncs.so
lib/syslog-ng/libconfgen.la
lib/syslog-ng/libconfgen.so
lib/syslog-ng/libconvertfuncs.la
lib/syslog-ng/libconvertfuncs.so
lib/syslog-ng/libcsvparser.la
lib/syslog-ng/libcsvparser.so
lib/syslog-ng/libdbparser.la
lib/syslog-ng/libdbparser.so
lib/syslog-ng/libdummy.la
lib/syslog-ng/libdummy.so
lib/syslog-ng/libsyslogformat.la
lib/syslog-ng/libsyslogformat.so
sbin/syslog-ng
sbin/syslog-ng-ctl
share/include/scl/pacct/plugin.conf
share/include/scl/syslogconf/README
share/include/scl/syslogconf/convert-syslogconf.awk
share/include/scl/syslogconf/plugin.conf
share/include/scl/system/generate-system-source.sh
share/include/scl/system/plugin.conf
share/xsd/patterndb-1.xsd
share/xsd/patterndb-2.xsd
share/xsd/patterndb-3.xsd
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%@dirrm %%DOCSDIR%%/sgml
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm share/xsd
@dirrm share/include/scl/system
@dirrm share/include/scl/syslogconf
@dirrm share/include/scl/pacct
@dirrm share/include/scl
@dirrm share/include
@dirrm lib/syslog-ng
@dirrm etc/patterndb.d
