# net-mgmt/ndpmon/files/ndpmonitor.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: ndpmonitor
# REQUIRE: DAEMON NETWORKING
#

. /etc/rc.subr

name="ndpmonitor"
rcvar=`set_rcvar`

command="%%PREFIX%%/sbin/ndpmon"
required_files=%%PREFIX%%/etc/ndpmon/config_ndpmon.xml

load_rc_config $name

: ${ndpmonitor_enable="NO"}

command_args=">/dev/null &"

run_rc_command "$1"
