# security/openvpn-auth-oauth2/pkg-plist
@sample %%ETCDIR%%/openvpn-auth-oauth2.yml.sample
etc/rc.d/openvpn_auth_oauth2
sbin/openvpn-auth-oauth2
