# sysutils/gkrellm2/pkg-plist
bin/gkrellm
bin/gkrellmd
include/gkrellm2/gkrellm.h
include/gkrellm2/gkrellm-public-proto.h
%%PORTDOCS%%share/doc/gkrellm/COPYRIGHT
%%PORTDOCS%%share/doc/gkrellm/CREDITS
%%PORTDOCS%%share/doc/gkrellm/README
%%PORTDOCS%%share/doc/gkrellm/Themes.html
share/locale/de/LC_MESSAGES/gkrellm.mo
share/locale/es/LC_MESSAGES/gkrellm.mo
share/locale/fr/LC_MESSAGES/gkrellm.mo
share/locale/ja/LC_MESSAGES/gkrellm.mo
share/locale/pl/LC_MESSAGES/gkrellm.mo
share/locale/ru/LC_MESSAGES/gkrellm.mo
share/locale/sl/LC_MESSAGES/gkrellm.mo
share/locale/sv/LC_MESSAGES/gkrellm.mo
%%PORTDOCS%%@dirrm share/doc/gkrellm
@unexec if cmp -s %D/etc/gkrellmd.conf %D/etc/gkrellmd.conf.sample; then rm -f %D/etc/gkrellmd.conf; fi
etc/gkrellmd.conf.sample
@exec [ -f %B/gkrellmd.conf ] || cp %B/%f %B/gkrellmd.conf
etc/rc.d/gkrellmd.sh.sample
@dirrm include/gkrellm2
@exec mkdir -p %D/libexec/gkrellm2/plugins
@exec mkdir -p %D/share/gkrellm2/themes
@unexec rmdir %D/libexec/gkrellm2/plugins 2>/dev/null || true
@unexec rmdir %D/libexec/gkrellm2 2>/dev/null || true
@unexec rmdir %D/share/gkrellm2/themes 2>/dev/null || true
@unexec rmdir %D/share/gkrellm2 2>/dev/null || true
