# databases/mysql323-server/pkg-plist
@unexec %D/etc/rc.d/mysql-server.sh stop > /dev/null 2>&1 || true
bin/comp_err
bin/isamchk
bin/isamlog
bin/my_print_defaults
bin/myisamchk
bin/myisamlog
bin/myisampack
bin/mysql_install_db
bin/pack_isam
bin/perror
bin/replace
bin/resolve_stack_dump
bin/resolveip
bin/safe_mysqld
lib/mysql/libdbug.a
lib/mysql/libheap.a
lib/mysql/libmerge.a
lib/mysql/libmyisam.a
lib/mysql/libmyisammrg.a
lib/mysql/libmystrings.a
lib/mysql/libmysys.a
lib/mysql/libnisam.a
libexec/mysqld
%%DATADIR%%/binary-configure
%%DATADIR%%/charsets/Index
%%DATADIR%%/charsets/README
%%DATADIR%%/charsets/cp1251.conf
%%DATADIR%%/charsets/cp1257.conf
%%DATADIR%%/charsets/croat.conf
%%DATADIR%%/charsets/danish.conf
%%DATADIR%%/charsets/dec8.conf
%%DATADIR%%/charsets/dos.conf
%%DATADIR%%/charsets/estonia.conf
%%DATADIR%%/charsets/german1.conf
%%DATADIR%%/charsets/greek.conf
%%DATADIR%%/charsets/hebrew.conf
%%DATADIR%%/charsets/hp8.conf
%%DATADIR%%/charsets/hungarian.conf
%%DATADIR%%/charsets/koi8_ru.conf
%%DATADIR%%/charsets/koi8_ukr.conf
%%DATADIR%%/charsets/latin1.conf
%%DATADIR%%/charsets/latin2.conf
%%DATADIR%%/charsets/latin5.conf
%%DATADIR%%/charsets/swe7.conf
%%DATADIR%%/charsets/usa7.conf
%%DATADIR%%/charsets/win1250.conf
%%DATADIR%%/charsets/win1251.conf
%%DATADIR%%/charsets/win1251ukr.conf
%%DATADIR%%/czech/errmsg.sys
%%DATADIR%%/czech/errmsg.txt
%%DATADIR%%/danish/errmsg.sys
%%DATADIR%%/danish/errmsg.txt
%%DATADIR%%/dutch/errmsg.sys
%%DATADIR%%/dutch/errmsg.txt
%%DATADIR%%/english/errmsg.sys
%%DATADIR%%/english/errmsg.txt
%%DATADIR%%/estonian/errmsg.sys
%%DATADIR%%/estonian/errmsg.txt
%%DATADIR%%/french/errmsg.sys
%%DATADIR%%/french/errmsg.txt
%%DATADIR%%/german/errmsg.sys
%%DATADIR%%/german/errmsg.txt
%%DATADIR%%/greek/errmsg.sys
%%DATADIR%%/greek/errmsg.txt
%%DATADIR%%/hungarian/errmsg.sys
%%DATADIR%%/hungarian/errmsg.txt
%%DATADIR%%/italian/errmsg.sys
%%DATADIR%%/italian/errmsg.txt
%%DATADIR%%/japanese/errmsg.sys
%%DATADIR%%/japanese/errmsg.txt
%%DATADIR%%/korean/errmsg.sys
%%DATADIR%%/korean/errmsg.txt
%%DATADIR%%/mi_test_all
%%DATADIR%%/mi_test_all.res
%%DATADIR%%/my-huge.cnf
%%DATADIR%%/my-large.cnf
%%DATADIR%%/my-medium.cnf
%%DATADIR%%/my-small.cnf
%%DATADIR%%/mysql-%%VER%%.spec
%%DATADIR%%/mysql-log-rotate
%%DATADIR%%/mysql.server
%%DATADIR%%/norwegian-ny/errmsg.sys
%%DATADIR%%/norwegian-ny/errmsg.txt
%%DATADIR%%/norwegian/errmsg.sys
%%DATADIR%%/norwegian/errmsg.txt
%%DATADIR%%/polish/errmsg.sys
%%DATADIR%%/polish/errmsg.txt
%%DATADIR%%/portuguese/errmsg.sys
%%DATADIR%%/portuguese/errmsg.txt
%%DATADIR%%/romanian/errmsg.sys
%%DATADIR%%/romanian/errmsg.txt
%%DATADIR%%/russian/errmsg.sys
%%DATADIR%%/russian/errmsg.txt
%%DATADIR%%/slovak/errmsg.sys
%%DATADIR%%/slovak/errmsg.txt
%%DATADIR%%/spanish/errmsg.sys
%%DATADIR%%/spanish/errmsg.txt
%%DATADIR%%/swedish/errmsg.sys
%%DATADIR%%/swedish/errmsg.txt
%%DATADIR%%/ukrainian/errmsg.sys
%%DATADIR%%/ukrainian/errmsg.txt
@unexec rmdir %D/lib/mysql 2> /dev/null || true
@dirrm %%DATADIR%%/charsets
@dirrm %%DATADIR%%/czech
@dirrm %%DATADIR%%/danish
@dirrm %%DATADIR%%/dutch
@dirrm %%DATADIR%%/english
@dirrm %%DATADIR%%/estonian
@dirrm %%DATADIR%%/french
@dirrm %%DATADIR%%/german
@dirrm %%DATADIR%%/greek
@dirrm %%DATADIR%%/hungarian
@dirrm %%DATADIR%%/italian
@dirrm %%DATADIR%%/japanese
@dirrm %%DATADIR%%/korean
@dirrm %%DATADIR%%/norwegian
@dirrm %%DATADIR%%/norwegian-ny
@dirrm %%DATADIR%%/polish
@dirrm %%DATADIR%%/portuguese
@dirrm %%DATADIR%%/romanian
@dirrm %%DATADIR%%/russian
@dirrm %%DATADIR%%/slovak
@dirrm %%DATADIR%%/spanish
@dirrm %%DATADIR%%/swedish
@dirrm %%DATADIR%%/ukrainian
@dirrm %%DATADIR%%
