# net-im/zephyr/files/zhm.in
#!/bin/sh
#
#
# PROVIDE: zhm
# REQUIRE: named

. /etc/rc.subr

name="zhm"
rcvar="zhm_enable"
start_cmd="zhm_start"
stop_postcmd="zhm_postcmd"
start_precmd="zhm_precmd"
command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/${name}.pid"

load_rc_config "$name"
eval "${rcvar}=\${${rcvar}:-'NO'}"

zhm_zephyr_servers="${zhm_zephyr_servers:-}"

zhm_precmd()
{
	mkdir -p $(dirname ${pidfile})
}

zhm_start()
{
	${command} -f ${zhm_zephyr_servers}
}

zhm_postcmd()
{
	rm -f ${pidfile}
}

run_rc_command "$1"
