# security/ipguard/files/openrc-ipguard.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf to enable ipguard:
#
# ipguard_interface="iface"
# ipguard_flags="-u 300 -xz"

name=ipguard
command=%%PREFIX%%/sbin/${name}

ipguard_interface=${ipguard_interface:-"fxp0"}
ipguard_flags=${ipguard_flags:-"-u 300 -xz"}
command_args="${ipguard_interface} ${ipguard_flags}"
pidfile="/var/run/${name}_${ipguard_interface}.pid"

depend() {
		keyword -shutdown
}
