# dns/maradns/pkg-plist
@comment $FreeBSD$
bin/askmara
etc/mararc.sample
etc/rc.d/maradns.sh
etc/rc.d/zoneserver.sh
etc/maradns/example_csv1
@unexec rmdir %D/etc/maradns 2>/dev/null || true
sbin/getzone
sbin/maradns
sbin/zoneserver
%%PORTDOCS%%%%DOCSDIR%%/00README.FIRST
%%PORTDOCS%%%%DOCSDIR%%/en/QuickStart
%%PORTDOCS%%%%DOCSDIR%%/en/README
%%PORTDOCS%%%%DOCSDIR%%/en/examples/example_authoritative_mararc.txt
%%PORTDOCS%%%%DOCSDIR%%/en/examples/example_csv1
%%PORTDOCS%%%%DOCSDIR%%/en/examples/example_csv1.txt
%%PORTDOCS%%%%DOCSDIR%%/en/examples/example_full_mararc
%%PORTDOCS%%%%DOCSDIR%%/en/examples/example_mararc
%%PORTDOCS%%%%DOCSDIR%%/en/examples/example_mararc.orig
%%PORTDOCS%%%%DOCSDIR%%/en/examples/example_recursive_mararc.txt
%%PORTDOCS%%%%DOCSDIR%%/en/faq.html
%%PORTDOCS%%%%DOCSDIR%%/en/faq.txt
%%PORTDOCS%%%%DOCSDIR%%/en/misc/0README
%%PORTDOCS%%%%DOCSDIR%%/en/misc/changelog.html
%%PORTDOCS%%%%DOCSDIR%%/en/misc/changelog.txt
%%PORTDOCS%%%%DOCSDIR%%/en/misc/credits.txt
%%PORTDOCS%%%%DOCSDIR%%/en/misc/faq.html
%%PORTDOCS%%%%DOCSDIR%%/en/misc/faq.txt
%%PORTDOCS%%%%DOCSDIR%%/en/misc/files.txt
%%PORTDOCS%%%%DOCSDIR%%/en/misc/recursive-case-insensitive
%%PORTDOCS%%%%DOCSDIR%%/en/misc/spammers/azmalink.net
%%PORTDOCS%%%%DOCSDIR%%/en/misc/spammers/hiddenonline/google.search
%%PORTDOCS%%%%DOCSDIR%%/en/misc/spammers/hiddenonline/hiddenonline.spam
%%PORTDOCS%%%%DOCSDIR%%/en/misc/spammers/hiddenonline/whois.dimeit.com
%%PORTDOCS%%%%DOCSDIR%%/en/misc/spammers/hiddenonline/whois.hiddenonline.net
%%PORTDOCS%%%%DOCSDIR%%/en/tutorial/authoritative.html
%%PORTDOCS%%%%DOCSDIR%%/en/tutorial/compile.html
%%PORTDOCS%%%%DOCSDIR%%/en/tutorial/convert.html
%%PORTDOCS%%%%DOCSDIR%%/en/tutorial/dangling.html
%%PORTDOCS%%%%DOCSDIR%%/en/tutorial/glossary.html
%%PORTDOCS%%%%DOCSDIR%%/en/tutorial/man.csv1.html
%%PORTDOCS%%%%DOCSDIR%%/en/tutorial/man.maradns.html
%%PORTDOCS%%%%DOCSDIR%%/en/tutorial/man.mararc.html
%%PORTDOCS%%%%DOCSDIR%%/en/tutorial/recordtypes.html
%%PORTDOCS%%%%DOCSDIR%%/en/tutorial/recursive.html
%%PORTDOCS%%%%DOCSDIR%%/en/tutorial/tutorial.html
%%PORTDOCS%%%%DOCSDIR%%/fr/QuickStart
%%PORTDOCS%%%%DOCSDIR%%/fr/faq.html
%%PORTDOCS%%%%DOCSDIR%%/fr/faq.txt
%%PORTDOCS%%%%DOCSDIR%%/fr/tutorial/authoritative.html
%%PORTDOCS%%%%DOCSDIR%%/fr/tutorial/compile.html
%%PORTDOCS%%%%DOCSDIR%%/fr/tutorial/convert.html
%%PORTDOCS%%%%DOCSDIR%%/fr/tutorial/glossary.html
%%PORTDOCS%%%%DOCSDIR%%/fr/tutorial/man.askmara.html
%%PORTDOCS%%%%DOCSDIR%%/fr/tutorial/man.csv1.html
%%PORTDOCS%%%%DOCSDIR%%/fr/tutorial/man.getzone.html
%%PORTDOCS%%%%DOCSDIR%%/fr/tutorial/man.maradns.html
%%PORTDOCS%%%%DOCSDIR%%/fr/tutorial/man.mararc.html
%%PORTDOCS%%%%DOCSDIR%%/fr/tutorial/man.zoneserver.html
%%PORTDOCS%%%%DOCSDIR%%/fr/tutorial/recursive.html
%%PORTDOCS%%%%DOCSDIR%%/fr/tutorial/tutorial.html
%%PORTDOCS%%@dirrm %%DOCSDIR%%/fr/tutorial
%%PORTDOCS%%@dirrm %%DOCSDIR%%/fr
%%PORTDOCS%%@dirrm %%DOCSDIR%%/en/tutorial
%%PORTDOCS%%@dirrm %%DOCSDIR%%/en/misc/spammers/hiddenonline
%%PORTDOCS%%@dirrm %%DOCSDIR%%/en/misc/spammers
%%PORTDOCS%%@dirrm %%DOCSDIR%%/en/misc
%%PORTDOCS%%@dirrm %%DOCSDIR%%/en/examples
%%PORTDOCS%%@dirrm %%DOCSDIR%%/en
%%PORTDOCS%%@dirrm %%DOCSDIR%%
