# x11/cde/files/dtspc.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: dtspc
# REQUIRE: rpcbind
# KEYWORD: shutdown

# Add the following line to /etc/rc.conf to enable dtspc:
#
# dtspc_enable=YES

. /etc/rc.subr

name=dtspc
rcvar=dtspc_enable
desc="CDE Subprocess Control Service"

load_rc_config $name

: ${dtspc_enable:=NO}

start_precmd="force_depend rpcbind || exit 1 && ln -sf %%PREFIX%%/dt /usr/dt && ln -sf %%PREFIX%%/etc/dt /etc/dt"
stop_postcmd="[ -L /usr/dt -a -L /etc/dt ] && rm /usr/dt /etc/dt"
pidfile=/var/run/dtspc.pid
command=/usr/sbin/daemon
command_args="-P $pidfile %%PREFIX%%/dt/bin/dtspcd"

run_rc_command $1
