# dns/unbound/files/unbound.in
#!/bin/sh
#
# unbound freebsd startup rc.d script, modified from the named script.
# uses the default unbound installation path and pidfile location.
# copy this to /etc/rc.d/unbound
# and put unbound_enable="YES" into rc.conf
#

# PROVIDE: unbound
# REQUIRE: SERVERS cleanvar
# KEYWORD: shutdown

. /etc/rc.subr

name="unbound"
rcvar=unbound_enable

command="%%PREFIX%%/sbin/unbound"
extra_commands="reload"
start_precmd="start_unbound_anchor"

load_rc_config $name

pidfile=${unbound_pidfile:-"%%PREFIX%%/etc/unbound/unbound.pid"}
unbound_enable=${unbound_enable:-"NO"}

start_unbound_anchor()
{
	echo -n "Obtaining a trust anchor:"
	%%PREFIX%%/sbin/unbound-anchor
	echo .
}

run_rc_command "$1"
