# sysutils/vm-bhyve-devel/pkg-plist
etc/rc.d/vm
etc/rc.d/vm_network
lib/vm-bhyve/vm-base
lib/vm-bhyve/vm-cmd
lib/vm-bhyve/vm-config
lib/vm-bhyve/vm-core
lib/vm-bhyve/vm-datastore
lib/vm-bhyve/vm-guest
lib/vm-bhyve/vm-info
lib/vm-bhyve/vm-migration
lib/vm-bhyve/vm-rctl
lib/vm-bhyve/vm-run
lib/vm-bhyve/vm-switch
lib/vm-bhyve/vm-switch-manual
lib/vm-bhyve/vm-switch-netgraph
lib/vm-bhyve/vm-switch-standard
lib/vm-bhyve/vm-switch-vale
lib/vm-bhyve/vm-switch-vxlan
lib/vm-bhyve/vm-util
lib/vm-bhyve/vm-zfs
share/man/man8/vm.8.gz
sbin/vm
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/config.sample
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/default.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/dragonfly.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/freebsd-zvol.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/freepbx.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/legacy-alpine.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/legacy-arch.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/legacy-debian.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/legacy-gentoo.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/legacy-ubuntu.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/linux.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/linux-grub.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/linux-grub-zvol.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/linux-zvol.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/netbsd.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/openbsd.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/resflash.conf
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/windows.conf
