# security/i2pd/files/i2pd.in
#!/bin/sh
#
# PROVIDE: i2pd
# REQUIRE: DAEMON NETWORKING FILESYSTEMS
# BEFORE: LOGIN
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable i2pd:
#
# i2pd_enable="YES"
#

. /etc/rc.subr

name="i2pd"
rcvar=i2pd_enable

load_rc_config ${name}

: ${i2pd_enable="NO"}
: ${i2pd_user="%%USER%%"}
: ${i2pd_group="%%GROUP%%"}
: ${i2pd_datadir="/var/db/${name}"}
: ${i2pd_config="%%ETCDIR%%/i2pd.conf"}

required_dirs=${i2pd_datadir}
pidfile=/var/run/i2pd/i2pd.pid
command="%%PREFIX%%/bin/${name}"
command_args="--conf ${i2pd_config} --service --datadir ${i2pd_datadir} --daemon"
extra_commands="reload"
reload_cmd="pkill -2 -F ${pidfile} -x i2pd"

run_rc_command "$1"
