# security/sguil/files/sguild.in
#!/bin/sh

# $FreeBSD: head/security/sguil-server/files/sguild.in 312467 2013-02-18 01:01:11Z miwi $

# PROVIDE: sguild
# REQUIRE: DAEMON
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable sguild:
# sguild_enable (bool):		Set to YES to enable sguild
# 				Default: NO
# sguild_flags (str):		Extra flags passed to sguild
#				Default: -D -P ${pid}
# sguild_conf (str):		Sguild configuration file
#				Default: %%PREFIX%%/etc/%%SGUILDIR%%/sguild.conf
# sguild_user (str):		Default: sguil
# Note: this value MUST be set in /etc/rc.conf if you do not accept the default
# user created by the pkg-install script

. /etc/rc.subr

name="sguild"
rcvar=sguild_enable
load_rc_config sguild
# set some defaults
: ${sguild_enable:="NO"}
: ${sguild_conf:="%%PREFIX%%/etc/%%SGUILDIR%%/sguild.conf"}
: ${pid:="/var/run/%%SGUILDIR%%/sguild.pid"}
: ${sguild_flags:="-D -P ${pid}"}
: ${sguild_user:="sguil"}

command="%%PREFIX%%/bin/${name}"
procname="%%PREFIX%%/bin/%%TCLSH%%"

run_rc_command "$1"
