# net/freevrrpd/pkg-plist
@sample etc/freevrrpd.conf.sample
etc/rc.d/freevrrpd
share/man/man8/freevrrpd.8.gz
sbin/freevrrpd
