# net/freevrrpd/pkg-plist
@unexec if cmp -s %D/etc/freevrrpd.conf %D/etc/freevrrpd.conf.sample; then rm -f %D/etc/freevrrpd.conf; fi
etc/freevrrpd.conf.sample
etc/rc.d/freevrrpd.sh.sample
sbin/freevrrpd
