# multimedia/freevo/files/freevo-rcsubr.in
#!/bin/sh

# PROVIDE: freevo
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Uses the following /etc/rc.conf variables:
#	freevo_enable="YES"			master on/off
#	freevo_recordserver_enable="YES"	start freevo record server
#	freevo_webserver_enable="YES"		start freevo web server
#	freevo_gui_enable="YES"			start freevo GUI
#

. /etc/rc.subr

name=freevo
rcvar=`set_rcvar`

load_rc_config $name
: ${freevo_enable="NO"}
: ${freevo_recordserver_enable="NO"}
: ${freevo_webserver_enable="NO"}
: ${freevo_gui_enable="NO"}

freevo_user=${freevo_user:-"freevo"}
freevo_group=${freevo_group:-"freevo"}
freevo_uid=`pw usershow ${freevo_user} | awk -F: '{print $3}'`

command=%%PREFIX%%/bin/freevo
command_interpreter=%%PREFIX%%/bin/python2.5
_interpbn=python

case "$1" in
fast*)	freevo_cmd=`expr "$1" : 'fast\(.*\)'` ;;
*)	freevo_cmd="$1" ;;
esac

name=freevo_recordserver
rcvar=${name}_enable
eval ${name}_user=${freevo_user}
eval ${name}_group=${freevo_group}
pidfile=/var/run/freevo/recordserver-${freevo_uid}.pid
command_args="recordserver ${freevo_cmd}"
run_rc_command "$1"

name=freevo_webserver
rcvar=${name}_enable
eval ${name}_user=${freevo_user}
eval ${name}_group=${freevo_group}
pidfile=/var/run/freevo/webserver-${freevo_uid}.pid
command_args="webserver ${freevo_cmd}"
run_rc_command "$1"

name=freevo_gui
rcvar=${name}_enable
eval ${name}_user=${freevo_user}
eval ${name}_group=${freevo_group}
pidfile=/var/run/freevo/freevo-${freevo_uid}.pid
command_args="${freevo_cmd}"
run_rc_command "$1"

