# dns/yadifa/files/openrc-yadifad.in
#!/sbin/openrc-run
#
# # optional
# yadifad_config="%%PREFIX%%/etc/yadifad.conf"

name=yadifad

: ${yadifad_config="%%PREFIX%%/etc/yadifad.conf"}

command=%%PREFIX%%/sbin/${name}
pidfile=/var/run/${name}.pid

command_args="-c ${yadifad_config} -d"
required_files=${yadifad_config}

depend() {
	keyword -shutdown
}
