# dns/dnrd/pkg-plist
@comment $FreeBSD$
etc/rc.d/dnrd.sh
sbin/dnrd
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README-cache
%%PORTDOCS%%%%DOCSDIR%%/README-master
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%@dirrm %%DOCSDIR%%
%%EXAMPLESDIR%%/master.sample
@dirrm %%EXAMPLESDIR%%
