# net/frickin/files/frickin.in
#!/bin/sh
#
# Startup script for frickin PPTP proxy.
#
# $FreeBSD$
#

# REQUIRE: DAEMON

# Define these frickin_* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#       /etc/rc.conf.d/frickin
#
# frickin_enable : bool
#   Enable the transport ("YES") or not ("NO", the default).
# frickin_config : string
#   Path to config file
#

frickin_enable=${frickin_enable-"NO"}
frickin_conf=${frickin_conf-""}

. /etc/rc.subr


name="frickin"
rcvar=frickin_enable

load_rc_config ${name}

command="%%PREFIX%%/sbin/${name}2"

if [ ! -z "${frickin_conf}" -a -f ${frickin_conf} ]; then
	frickin_flags="-c ${frickin_conf}"
fi

run_rc_command "$1"
