# japanese/rskkserv/pkg-plist
etc/rc.d/rskkserv.sh
@unexec if cmp -s %D/etc/rskkserv.conf %D/%%EXAMPLESDIR%%/rskkserv.conf; then rm -f %D/etc/rskkserv.conf; fi
%%RUBY_SITEARCHLIBDIR%%/skkserv/skkdic.so
@dirrm %%RUBY_SITEARCHLIBDIR%%/skkserv
%%RUBY_SITELIBDIR%%/skkserv/bdbdic.rb
%%RUBY_SITELIBDIR%%/skkserv/conf.rb
%%RUBY_SITELIBDIR%%/skkserv/ebdic.rb
%%RUBY_SITELIBDIR%%/skkserv/logger.rb
%%RUBY_SITELIBDIR%%/skkserv/skkdic.rb
@dirrm %%RUBY_SITELIBDIR%%/skkserv
sbin/rskkserv
@unexec rm -f %D/share/skk/SKK-JISYO.[LMS].[an] 2>/dev/null || true
%%PORTDOCS%%%%RUBY_MODDOCDIR%%/ChangeLog
%%PORTDOCS%%%%RUBY_MODDOCDIR%%/ja/TODO
%%PORTDOCS%%@dirrm %%RUBY_MODDOCDIR%%/ja
%%PORTDOCS%%@dirrm %%RUBY_MODDOCDIR%%
%%EXAMPLESDIR%%/conf-o2n.rb
%%EXAMPLESDIR%%/rskkserv.conf
@exec [ -f %D/etc/rskkserv.conf ] || cp %B/%f %D/etc/rskkserv.conf
%%EXAMPLESDIR%%/rskkserv.conf.sample
@dirrm %%EXAMPLESDIR%%
