# databases/monetdb/files/openrc-monetdbd.in
#!/sbin/openrc-run

name=monetdbd
command=%%PREFIX%%/bin/monetdbd
required_dirs=%%PREFIX%%/monetdb5/dbfarm
pidfile=/var/run/monetdb/merovingian.pid

start() {
	su -l monetdb -c "exec %%PREFIX%%/bin/monetdbd start %%PREFIX%%/monetdb5/dbfarm"
}

start_post() {
	%%PREFIX%%/bin/monetdb start -a
}

depend() {
	keyword -shutdown
}
