# net/haproxy/pkg-plist
sbin/haproxy
etc/rc.d/haproxy.sh
%%PORTDOCS%%%%DOCSDIR%%/architecture.txt
%%PORTDOCS%%%%DOCSDIR%%/haproxy-en.txt
%%PORTDOCS%%%%DOCSDIR%%/haproxy-fr.txt
%%PORTDOCS%%%%EXAMPLESDIR%%/haproxy.cfg
%%PORTDOCS%%@dirrm %%DOCSDIR%%
%%PORTDOCS%%@dirrm %%EXAMPLESDIR%%
