# security/osiris/pkg-plist
@comment $FreeBSD$
%%CLI%%sbin/osiris
%%OSIRISD%%sbin/osirisd
%%OSIRISMD%%sbin/osirismd
%%OSIRISD%%etc/rc.d/osirisd.sh
%%OSIRISMD%%etc/rc.d/osirismd.sh
@exec mkdir -p %D/osiris
@exec chown osiris:osiris %D/osiris
%%OSIRISMD%%osiris/configs/default.aix
%%OSIRISMD%%osiris/configs/default.bsdos
%%OSIRISMD%%osiris/configs/default.darwin
%%OSIRISMD%%osiris/configs/default.freebsd
%%OSIRISMD%%osiris/configs/default.irix
%%OSIRISMD%%osiris/configs/default.linux
%%OSIRISMD%%osiris/configs/default.netbsd
%%OSIRISMD%%osiris/configs/default.openbsd
%%OSIRISMD%%osiris/configs/default.sunos
%%OSIRISMD%%osiris/configs/default.unix-generic
%%OSIRISMD%%osiris/configs/default.windows2000
%%OSIRISMD%%osiris/configs/default.windowsnt
%%OSIRISMD%%osiris/configs/default.windowsserver2003
%%OSIRISMD%%osiris/configs/default.windowsxp
%%OSIRISMD%%osiris/configs/default.wrt54g
%%OSIRISMD%%@dirrm osiris/configs
@dirrm osiris
