# security/osiris/pkg-plist
%%CLI%%sbin/osiris
%%PRINTDB%%bin/printdb
%%OSIRISD%%sbin/osirisd
%%OSIRISMD%%@owner osiris
%%OSIRISMD%%@group osiris
%%OSIRISMD%%sbin/osirismd
%%OSIRISMD%%@owner
%%OSIRISMD%%@group
%%OSIRISD%%etc/rc.d/osirisd.sh
%%OSIRISMD%%etc/rc.d/osirismd.sh
@dir(osiris,osiris,) osiris
%%OSIRISMD%%@owner osiris
%%OSIRISMD%%@group osiris
%%OSIRISMD%%osiris/configs/default.aix
%%OSIRISMD%%osiris/configs/default.bsdos
%%OSIRISMD%%osiris/configs/default.darwin
%%OSIRISMD%%osiris/configs/default.freebsd
%%OSIRISMD%%osiris/configs/default.irix
%%OSIRISMD%%osiris/configs/default.linux
%%OSIRISMD%%osiris/configs/default.netbsd
%%OSIRISMD%%osiris/configs/default.openbsd
%%OSIRISMD%%osiris/configs/default.sunos
%%OSIRISMD%%osiris/configs/default.unix-generic
%%OSIRISMD%%osiris/configs/default.windows2000
%%OSIRISMD%%osiris/configs/default.windowsnt
%%OSIRISMD%%osiris/configs/default.windowsserver2003
%%OSIRISMD%%osiris/configs/default.windowsxp
%%OSIRISMD%%osiris/configs/default.wrt54g
%%OSIRISMD%%@owner
%%OSIRISMD%%@group
%%OSIRISMD%%@dirrm osiris/configs
