# net-im/prosody/files/prosody.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: prosody
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# prosody_enable (bool):   Set to NO by default.
#                Set it to YES to enable Prosody

. /etc/rc.subr

name="prosody"
rcvar="prosody_enable"
command="%%PREFIX%%/bin/prosodyctl"

load_rc_config $name

prosody_enable=${prosody_enable:-"NO"}
pidfile=${prosody_pidfile:-"/var/run/prosody/prosody.pid"}

extra_commands="reload"
start_cmd="prosody_cmd start"
stop_cmd="prosody_cmd stop"
reload_cmd="prosody_cmd reload"
restart_cmd="prosody_cmd restart"
status_cmd="prosody_cmd status"

export PATH="%%PREFIX%%/bin:$PATH"

prosody_cmd()
{
	$command $1
}

run_rc_command "$1"
