# sysutils/mfid/files/mfid.in
#!/bin/sh
#
# $FreeBSD$

# PROVIDE: mfid
# REQUIRE: root ldconfig devfs syslogd mail
# KEYWORD: nojail noyroot shutdown

. /etc/rc.subr

name=mfid
rcvar=mfid_enable
load_rc_config $name

required_files=/dev/mfi0

raid_alert_mailto=${mfid_alert_mailto:-"root@localhost"}
command=%%PREFIX%%/sbin/mfid
command_args=${raid_alert_mailto}

run_rc_command "$1"
