# irc/irc/pkg-plist
bin/irc
etc/rc.d/ircd.sh
etc/ircd/iauth.conf
etc/ircd/ircd.m4
etc/ircd/example.conf
sbin/chkconf
sbin/ircd
sbin/iauth
sbin/ircd-mkpasswd
sbin/ircdwatch
@dirrm etc/ircd
@dirrm var/log/ircd
@dirrm var/log
@dirrm var/run
@dirrm var
