# www/www6to4/pkg-plist
sbin/www6to4
etc/www6to4/www6to4.conf.sample
etc/www6to4/www6to4_forward.conf.sample
etc/rc.d/www6to4.sh
%%PORTDOCS%%share/doc/www6to4/README
@dirrm etc/www6to4
%%PORTDOCS%%@dirrm share/doc/www6to4
