# www/miniflux/files/miniflux.in
#!/bin/sh
#
# Author: Mark Felder <feld@FreeBSd.org>
#
# $FreeBSD$
#

# PROVIDE: miniflux
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable miniflux:
# miniflux_enable="YES"
# miniflux_flags="<set as needed>"

. /etc/rc.subr

name="miniflux"
rcvar=miniflux_enable

load_rc_config $name

: ${miniflux_enable="NO"}

# Miniflux doesn’t use any config file, only environment variables.
set -a
. %%PREFIX%%/etc/miniflux.env
set +a

pidfile="/var/run/miniflux.pid"
command=/usr/sbin/daemon
command_args=" -p ${pidfile} -u miniflux %%PREFIX%%/bin/miniflux | /usr/bin/logger -t miniflux &"

run_rc_command "$1"
