# www/miniflux/files/miniflux.in
#!/bin/sh
#
# Author: Mark Felder <feld@FreeBSd.org>
#
# $FreeBSD$
#

# PROVIDE: miniflux
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable miniflux:
# miniflux_enable="YES"
# miniflux_flags="<set as needed>"

. /etc/rc.subr

name="miniflux"
rcvar=miniflux_enable

load_rc_config $name

: ${miniflux_enable="NO"}
: ${miniflux_config="%%PREFIX%%/etc/miniflux.env"}

pidfile="/var/run/miniflux.pid"
procname="%%PREFIX%%/bin/miniflux"
command="/usr/sbin/daemon"
command_args="-f -T miniflux -p ${pidfile} -u miniflux %%PREFIX%%/bin/miniflux -c ${miniflux_config}"

run_rc_command "$1"
