# japanese/esecanna/pkg-plist
etc/esecannarc
etc/rc.d/%%STARTER%%
@exec mkdir -p %D/lib/esecanna
@dirrm lib/esecanna
sbin/esecannaserver
share/doc/esecanna/CHANGES.jp
share/doc/esecanna/README.First.jp
share/doc/esecanna/README.jp
@dirrm share/doc/esecanna
