# japanese/esecanna/pkg-plist
@unexec if [ -x %D/etc/rc.d/esecanna.sh ]; then %D/etc/rc.d/esecanna.sh stop || true; fi
etc/esecannarc
etc/rc.d/esecanna.sh
@exec mkdir -p %D/lib/esecanna
@dirrm lib/esecanna
sbin/esecannaserver
@exec /bin/mkdir -p /var/run/esecanna
@unexec /bin/rm -rf /var/run/esecanna
%%PORTDOCS%%share/doc/esecanna/CREDITS
%%PORTDOCS%%share/doc/esecanna/CHANGES.jp
%%PORTDOCS%%share/doc/esecanna/README.First.jp
%%PORTDOCS%%share/doc/esecanna/README.jp
%%PORTDOCS%%@dirrm share/doc/esecanna
