# japanese/esecanna/pkg-plist
@unexec if [ -x %D/etc/rc.d/esecanna.sh ]; then %D/etc/rc.d/esecanna.sh stop || true; fi
etc/esecannarc
etc/rc.d/%%STARTER%%
@unexec if [ -f %B/canna.sh ]; then chmod +x %B/canna.sh; fi
@exec if [ -x %B/canna.sh ]; then chmod -x %B/canna.sh; fi
@exec mkdir -p %D/lib/esecanna
@dirrm lib/esecanna
sbin/esecannaserver
%%PORTDOCS%%share/doc/esecanna/CREDITS
%%PORTDOCS%%share/doc/esecanna/CHANGES.jp
%%PORTDOCS%%share/doc/esecanna/README.First.jp
%%PORTDOCS%%share/doc/esecanna/README.jp
%%PORTDOCS%%@dirrm share/doc/esecanna
