# net/libcharon/pkg-plist
%%PYTHON_SITELIBDIR%%/Charon/Client/DBusInterface.py
%%PYTHON_SITELIBDIR%%/Charon/Client/Request.py
%%PYTHON_SITELIBDIR%%/Charon/Client/__init__.py
%%PYTHON_SITELIBDIR%%/Charon/Client/test_glib.py
%%PYTHON_SITELIBDIR%%/Charon/Client/test_qt.py
%%PYTHON_SITELIBDIR%%/Charon/FileInterface.py
%%PYTHON_SITELIBDIR%%/Charon/OpenMode.py
%%PYTHON_SITELIBDIR%%/Charon/ReadOnlyError.py
%%PYTHON_SITELIBDIR%%/Charon/Service/FileService.py
%%PYTHON_SITELIBDIR%%/Charon/Service/RequestQueue.py
%%PYTHON_SITELIBDIR%%/Charon/Service/__init__.py
%%PYTHON_SITELIBDIR%%/Charon/Service/main.py
%%PYTHON_SITELIBDIR%%/Charon/VirtualFile.py
%%PYTHON_SITELIBDIR%%/Charon/WriteOnlyError.py
%%PYTHON_SITELIBDIR%%/Charon/__init__.py
%%PYTHON_SITELIBDIR%%/Charon/filetypes/GCodeFile.py
%%PYTHON_SITELIBDIR%%/Charon/filetypes/GCodeGzFile.py
%%PYTHON_SITELIBDIR%%/Charon/filetypes/OpenPackagingConvention.py
%%PYTHON_SITELIBDIR%%/Charon/filetypes/UltimakerFormatPackage.py
%%PYTHON_SITELIBDIR%%/Charon/filetypes/__init__.py
lib/systemd/system/charon.service
share/dbus-1/system.d/nl.ultimaker.charon.conf
