# dns/dnscrypt-proxy/files/dnscrypt-proxy.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: dnscrypt-proxy
# REQUIRE: SERVERS cleanvar
# BEFORE: named
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable dnscrypt-proxy:
#
# dnscrypt_proxy_enable="YES":	Set to NO by default.
#				Set it to YES to enable dnscrypt-proxy.
#
# To redirect a local resolver through dnscrypt-proxy, point it at 127.0.0.2
# and add the following to rc.conf:
# ifconfig_lo0_alias0="inet 127.0.0.2 netmask 0xffffffff"
# dnscrypt_proxy_flags='-a 127.0.0.2'

. /etc/rc.subr

name=dnscrypt_proxy
rcvar=dnscrypt_proxy_enable

stop_cmd=dnscrypt_proxy_stop

load_rc_config ${name}

: ${dnscrypt_proxy_enable:=NO}
: ${dnscrypt_proxy_uid=_dnscrypt-proxy} # User to run daemon as
: ${dnscrypt_proxy_pidfile=/var/run/dnscrypt-proxy.pid} # Path to pid file
: ${dnscrypt_proxy_logfile=/var/log/dnscrypt-proxy.log} # Path to log file

if [ -n "$dnscrypt_proxy_uid" ]; then
	dnscrypt_proxy_flags="${dnscrypt_proxy_flags} -u ${dnscrypt_proxy_uid}"
fi

command=%%PREFIX%%/sbin/dnscrypt-proxy
procname=%%PREFIX%%/sbin/dnscrypt-proxy

command_args="-d -p ${dnscrypt_proxy_pidfile} -l ${dnscrypt_proxy_logfile}"

dnscrypt_proxy_stop() {
        kill -KILL `cat ${pidfile}` 2> /dev/null && echo "Killed ${name}."
        }

run_rc_command "$1"
