# sysutils/dockerbox/pkg-plist
@sample etc/dockerbox/dockerbox.conf.sample
etc/rc.d/dockerbox
sbin/dockerbox
