# russian/nagios/pkg-plist
@exec mkdir -p %D/%%DATADIR%%//ssi
@exec mkdir -p /var/spool/nagios/archives
@exec mkdir -p /var/spool/nagios/rw
@exec chmod 775 /var/spool/nagios /var/spool/nagios/archives /var/spool/nagios/rw
@exec chown nagios:nagios /var/spool/nagios /var/spool/nagios/archives
@exec chown nagios:nogroup /var/spool/nagios/rw
bin/nagios
etc/nagios/cgi.cfg-sample
etc/nagios/checkcommands.cfg-sample
etc/nagios/contactgroups.cfg-sample
etc/nagios/contacts.cfg-sample
etc/nagios/dependencies.cfg-sample
etc/nagios/escalations.cfg-sample
etc/nagios/hostextinfo.cfg-sample
etc/nagios/hostgroups.cfg-sample
etc/nagios/hosts.cfg-sample
etc/nagios/misccommands.cfg-sample
etc/nagios/nagios.cfg-sample
etc/nagios/resource.cfg-sample
etc/nagios/serviceextinfo.cfg-sample
etc/nagios/services.cfg-sample
etc/nagios/timeperiods.cfg-sample
etc/rc.d/nagios.sh
%%DATADIR%%/cgi-bin/avail.cgi
%%DATADIR%%/cgi-bin/cmd.cgi
%%DATADIR%%/cgi-bin/config.cgi
%%DATADIR%%/cgi-bin/extinfo.cgi
%%DATADIR%%/cgi-bin/histogram.cgi
%%DATADIR%%/cgi-bin/history.cgi
%%DATADIR%%/cgi-bin/notifications.cgi
%%DATADIR%%/cgi-bin/outages.cgi
%%DATADIR%%/cgi-bin/showlog.cgi
%%DATADIR%%/cgi-bin/status.cgi
%%DATADIR%%/cgi-bin/statusmap.cgi
%%DATADIR%%/cgi-bin/statuswml.cgi
%%DATADIR%%/cgi-bin/statuswrl.cgi
%%DATADIR%%/cgi-bin/summary.cgi
%%DATADIR%%/cgi-bin/tac.cgi
%%DATADIR%%/cgi-bin/trends.cgi
%%DATADIR%%/contexthelp/A1.html
%%DATADIR%%/contexthelp/A2.html
%%DATADIR%%/contexthelp/A3.html
%%DATADIR%%/contexthelp/A4.html
%%DATADIR%%/contexthelp/A5.html
%%DATADIR%%/contexthelp/A6.html
%%DATADIR%%/contexthelp/A7.html
%%DATADIR%%/contexthelp/B1.html
%%DATADIR%%/contexthelp/C1.html
%%DATADIR%%/contexthelp/D1.html
%%DATADIR%%/contexthelp/E1.html
%%DATADIR%%/contexthelp/F1.html
%%DATADIR%%/contexthelp/G1.html
%%DATADIR%%/contexthelp/G2.html
%%DATADIR%%/contexthelp/G3.html
%%DATADIR%%/contexthelp/G4.html
%%DATADIR%%/contexthelp/G5.html
%%DATADIR%%/contexthelp/G6.html
%%DATADIR%%/contexthelp/H1.html
%%DATADIR%%/contexthelp/H2.html
%%DATADIR%%/contexthelp/H3.html
%%DATADIR%%/contexthelp/H4.html
%%DATADIR%%/contexthelp/H5.html
%%DATADIR%%/contexthelp/H6.html
%%DATADIR%%/contexthelp/H7.html
%%DATADIR%%/contexthelp/H8.html
%%DATADIR%%/contexthelp/I1.html
%%DATADIR%%/contexthelp/I2.html
%%DATADIR%%/contexthelp/I3.html
%%DATADIR%%/contexthelp/I4.html
%%DATADIR%%/contexthelp/I5.html
%%DATADIR%%/contexthelp/I6.html
%%DATADIR%%/contexthelp/I7.html
%%DATADIR%%/contexthelp/I8.html
%%DATADIR%%/contexthelp/J1.html
%%DATADIR%%/contexthelp/K1.html
%%DATADIR%%/contexthelp/L1.html
%%DATADIR%%/contexthelp/L10.html
%%DATADIR%%/contexthelp/L11.html
%%DATADIR%%/contexthelp/L12.html
%%DATADIR%%/contexthelp/L13.html
%%DATADIR%%/contexthelp/L2.html
%%DATADIR%%/contexthelp/L3.html
%%DATADIR%%/contexthelp/L4.html
%%DATADIR%%/contexthelp/L5.html
%%DATADIR%%/contexthelp/L6.html
%%DATADIR%%/contexthelp/L7.html
%%DATADIR%%/contexthelp/L8.html
%%DATADIR%%/contexthelp/L9.html
%%DATADIR%%/contexthelp/M1.html
%%DATADIR%%/contexthelp/M2.html
%%DATADIR%%/contexthelp/M3.html
%%DATADIR%%/contexthelp/M4.html
%%DATADIR%%/contexthelp/M5.html
%%DATADIR%%/contexthelp/M6.html
%%DATADIR%%/contexthelp/N1.html
%%DATADIR%%/contexthelp/N2.html
%%DATADIR%%/contexthelp/N3.html
%%DATADIR%%/contexthelp/N4.html
%%DATADIR%%/contexthelp/N5.html
%%DATADIR%%/contexthelp/N6.html
%%DATADIR%%/contexthelp/N7.html
%%DATADIR%%/docs/images/activepassive.png
%%DATADIR%%/docs/images/cgi-avail-a.png
%%DATADIR%%/docs/images/cgi-avail-b.png
%%DATADIR%%/docs/images/cgi-cmd.png
%%DATADIR%%/docs/images/cgi-config.png
%%DATADIR%%/docs/images/cgi-extinfo-a.png
%%DATADIR%%/docs/images/cgi-extinfo-b.png
%%DATADIR%%/docs/images/cgi-extinfo-c.png
%%DATADIR%%/docs/images/cgi-extinfo-d.png
%%DATADIR%%/docs/images/cgi-histogram.png
%%DATADIR%%/docs/images/cgi-history.png
%%DATADIR%%/docs/images/cgi-notifications.png
%%DATADIR%%/docs/images/cgi-outages.png
%%DATADIR%%/docs/images/cgi-showlog.png
%%DATADIR%%/docs/images/cgi-status-a.png
%%DATADIR%%/docs/images/cgi-status-b.png
%%DATADIR%%/docs/images/cgi-status-c.png
%%DATADIR%%/docs/images/cgi-status-d.png
%%DATADIR%%/docs/images/cgi-statusmap.png
%%DATADIR%%/docs/images/cgi-statuswml.png
%%DATADIR%%/docs/images/cgi-statuswrl.png
%%DATADIR%%/docs/images/cgi-summary.png
%%DATADIR%%/docs/images/cgi-tac.png
%%DATADIR%%/docs/images/cgi-trends.png
%%DATADIR%%/docs/images/checktiming.png
%%DATADIR%%/docs/images/distributed.png
%%DATADIR%%/docs/images/host-dependencies.png
%%DATADIR%%/docs/images/indirecthostcheck.png
%%DATADIR%%/docs/images/indirectsvccheck.png
%%DATADIR%%/docs/images/indirectsvccheck2.png
%%DATADIR%%/docs/images/interleaved1.png
%%DATADIR%%/docs/images/interleaved2.png
%%DATADIR%%/docs/images/interleaved3.png
%%DATADIR%%/docs/images/logofullsize.png
%%DATADIR%%/docs/images/network-heirarchy.png
%%DATADIR%%/docs/images/network-outage1.png
%%DATADIR%%/docs/images/network-outage2.png
%%DATADIR%%/docs/images/noninterleaved1.png
%%DATADIR%%/docs/images/noninterleaved2.png
%%DATADIR%%/docs/images/physical-network.png
%%DATADIR%%/docs/images/plugintheory.png
%%DATADIR%%/docs/images/redudancy.png
%%DATADIR%%/docs/images/redundancy.png
%%DATADIR%%/docs/images/service-dependencies.png
%%DATADIR%%/docs/images/statetransitions.png
%%DATADIR%%/docs/images/statetransitions2.png
%%DATADIR%%/docs/about.html
%%DATADIR%%/docs/addons.html
%%DATADIR%%/docs/beginners.html
%%DATADIR%%/docs/cgiauth.html
%%DATADIR%%/docs/cgiincludes.html
%%DATADIR%%/docs/cgis.html
%%DATADIR%%/docs/checkscheduling.html
%%DATADIR%%/docs/clusters.html
%%DATADIR%%/docs/commandfile.html
%%DATADIR%%/docs/config.html
%%DATADIR%%/docs/configcgi.html
%%DATADIR%%/docs/configextinfo.html
%%DATADIR%%/docs/configmain.html
%%DATADIR%%/docs/configobject.html
%%DATADIR%%/docs/dependencies.html
%%DATADIR%%/docs/distributed.html
%%DATADIR%%/docs/downtime.html
%%DATADIR%%/docs/embeddedperl.html
%%DATADIR%%/docs/escalations.html
%%DATADIR%%/docs/eventhandlers.html
%%DATADIR%%/docs/extcommands.html
%%DATADIR%%/docs/faqs.html
%%DATADIR%%/docs/flapping.html
%%DATADIR%%/docs/freshness.html
%%DATADIR%%/docs/funstuff.html
%%DATADIR%%/docs/index.html
%%DATADIR%%/docs/indirectchecks.html
%%DATADIR%%/docs/installing.html
%%DATADIR%%/docs/installweb.html
%%DATADIR%%/docs/int-portsentry.html
%%DATADIR%%/docs/int-snmptrap.html
%%DATADIR%%/docs/int-tcpwrappers.html
%%DATADIR%%/docs/macros.html
%%DATADIR%%/docs/networkoutages.html
%%DATADIR%%/docs/networkreachability.html
%%DATADIR%%/docs/notifications.html
%%DATADIR%%/docs/parallelization.html
%%DATADIR%%/docs/passivechecks.html
%%DATADIR%%/docs/perfdata.html
%%DATADIR%%/docs/plugins.html
%%DATADIR%%/docs/plugintheory.html
%%DATADIR%%/docs/redundancy.html
%%DATADIR%%/docs/robots.txt
%%DATADIR%%/docs/security.html
%%DATADIR%%/docs/stalking.html
%%DATADIR%%/docs/starting.html
%%DATADIR%%/docs/statetypes.html
%%DATADIR%%/docs/stoprestart.html
%%DATADIR%%/docs/templaterecursion.html
%%DATADIR%%/docs/templatetricks.html
%%DATADIR%%/docs/timeperiods.html
%%DATADIR%%/docs/toc.html
%%DATADIR%%/docs/tuning.html
%%DATADIR%%/docs/verifyconfig.html
%%DATADIR%%/docs/volatileservices.html
%%DATADIR%%/docs/whatsnew.html
%%DATADIR%%/docs/xdata-db.html
%%DATADIR%%/docs/xeddefault.html
%%DATADIR%%/docs/xedtemplate.html
%%DATADIR%%/docs/xoddefault.html
%%DATADIR%%/docs/xodtemplate.html
%%DATADIR%%/docs/xpddefault.html
%%DATADIR%%/docs/xpdfile.html
%%DATADIR%%/images/ack.gif
%%DATADIR%%/images/command.png
%%DATADIR%%/images/comment.gif
%%DATADIR%%/images/contexthelp1.gif
%%DATADIR%%/images/contexthelp2.gif
%%DATADIR%%/images/critical.png
%%DATADIR%%/images/delay.gif
%%DATADIR%%/images/delete.gif
%%DATADIR%%/images/disabled.gif
%%DATADIR%%/images/down.gif
%%DATADIR%%/images/downtime.gif
%%DATADIR%%/images/empty.gif
%%DATADIR%%/images/enabled.gif
%%DATADIR%%/images/extinfo.gif
%%DATADIR%%/images/flapping.gif
%%DATADIR%%/images/greendot.gif
%%DATADIR%%/images/histogram.png
%%DATADIR%%/images/history.gif
%%DATADIR%%/images/hostevent.gif
%%DATADIR%%/images/info.png
%%DATADIR%%/images/left.gif
%%DATADIR%%/images/logofullsize.jpg
%%DATADIR%%/images/logos/nagios.gd2
%%DATADIR%%/images/logos/nagios.gif
%%DATADIR%%/images/logos/nagiosvrml.png
%%DATADIR%%/images/logos/unknown.gd2
%%DATADIR%%/images/logos/unknown.gif
%%DATADIR%%/images/logrotate.png
%%DATADIR%%/images/ndisabled.gif
%%DATADIR%%/images/noack.gif
%%DATADIR%%/images/notes.gif
%%DATADIR%%/images/notify.gif
%%DATADIR%%/images/orangedot.gif
%%DATADIR%%/images/passiveonly.gif
%%DATADIR%%/images/recovery.png
%%DATADIR%%/images/redudancy.png
%%DATADIR%%/images/redundancy.png
%%DATADIR%%/images/restart.gif
%%DATADIR%%/images/right.gif
%%DATADIR%%/images/sbconfig.png
%%DATADIR%%/images/sbgeneral.png
%%DATADIR%%/images/sblogo.jpg
%%DATADIR%%/images/sbmonitor.png
%%DATADIR%%/images/sbreport.png
%%DATADIR%%/images/serviceevent.gif
%%DATADIR%%/images/start.gif
%%DATADIR%%/images/status.gif
%%DATADIR%%/images/status2.gif
%%DATADIR%%/images/status3.gif
%%DATADIR%%/images/status4.gif
%%DATADIR%%/images/stop.gif
%%DATADIR%%/images/tacdisabled.jpg
%%DATADIR%%/images/tacdisabled.png
%%DATADIR%%/images/tacenabled.jpg
%%DATADIR%%/images/tacenabled.png
%%DATADIR%%/images/thermcrit.png
%%DATADIR%%/images/thermok.png
%%DATADIR%%/images/thermwarn.png
%%DATADIR%%/images/trends.gif
%%DATADIR%%/images/trendshost.png
%%DATADIR%%/images/trendssvc.png
%%DATADIR%%/images/unknown.png
%%DATADIR%%/images/up.gif
%%DATADIR%%/images/warning.png
%%DATADIR%%/images/weblogo1.png
%%DATADIR%%/images/zoom1.gif
%%DATADIR%%/images/zoom2.gif
%%DATADIR%%/index.html
%%DATADIR%%/main.html
%%DATADIR%%/media/critical.wav
%%DATADIR%%/media/hostdown.wav
%%DATADIR%%/media/warning.wav
%%DATADIR%%/robots.txt
%%DATADIR%%/side.html
%%DATADIR%%/stylesheets/avail.css
%%DATADIR%%/stylesheets/checksanity.css
%%DATADIR%%/stylesheets/cmd.css
%%DATADIR%%/stylesheets/config.css
%%DATADIR%%/stylesheets/extinfo.css
%%DATADIR%%/stylesheets/histogram.css
%%DATADIR%%/stylesheets/history.css
%%DATADIR%%/stylesheets/ministatus.css
%%DATADIR%%/stylesheets/notifications.css
%%DATADIR%%/stylesheets/outages.css
%%DATADIR%%/stylesheets/showlog.css
%%DATADIR%%/stylesheets/status.css
%%DATADIR%%/stylesheets/statusmap.css
%%DATADIR%%/stylesheets/summary.css
%%DATADIR%%/stylesheets/tac.css
%%DATADIR%%/stylesheets/trends.css
@unexec rmdir %D/etc/nagios 2>/dev/null || true
@dirrm %%DATADIR%%/cgi-bin
@dirrm %%DATADIR%%/contexthelp
@dirrm %%DATADIR%%/docs/images
@dirrm %%DATADIR%%/docs
@dirrm %%DATADIR%%/images/logos
@dirrm %%DATADIR%%/images
@dirrm %%DATADIR%%/media
@dirrm %%DATADIR%%/ssi
@dirrm %%DATADIR%%/stylesheets
@dirrm %%DATADIR%%/
@unexec rmdir /var/spool/nagios/archives 2>/dev/null || true
@unexec rmdir /var/spool/nagios/rw 2>/dev/null || true
@unexec rmdir /var/spool/nagios 2>/dev/null || true
