# net-mgmt/pmacct/files/pmacctd.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: pmacctd
# REQUIRE: NETWORKING
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable pmacctd:
#
# pmacctd_enable="YES"
#
# You can define flags for pmacctd running. See pmacctd(8) for details.
#
# Example:
# pmacctd_flags="-f %%PREFIX%%/etc/pmacctd.conf -i em0"
#
# Default flags: -f %%PREFIX%%/etc/pmacctd.conf
#

. /etc/rc.subr

name=pmacctd
rcvar=pmacctd_enable

command=%%PREFIX%%/sbin/$name
command_args="-D"

# set defaults

pmacctd_enable=${pmacctd_enable:-"NO"}
pmacctd_flags=${pmacctd_flags:-"-f %%PREFIX%%/etc/pmacctd.conf"}

load_rc_config $name
run_rc_command "$1"
