# www/yarr/files/yarr.in
#!/bin/sh

# PROVIDE: yarr
# REQUIRE: NETWORKING

. /etc/rc.subr

name="yarr"
rcvar="yarr_enable"

pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args="-P ${pidfile} -r -f -o /var/log/${name}.log %%PREFIX%%/bin/${name}"

load_rc_config $name
: ${yarr_enable:=NO}

run_rc_command "$1"
