# www/yarr/files/yarr.in
#!/bin/sh

# PROVIDE: yarr
# REQUIRE: NETWORKING

# To set various configuration parameters it's recommended to set yarr_env
# in rc.conf.
#
# yarr_env="YARR_ADDR=... YARR_BASE=..."

. /etc/rc.subr

name="yarr"
rcvar="yarr_enable"

pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args="-P ${pidfile} -r -f -o /var/log/${name}.log %%PREFIX%%/bin/${name}"

load_rc_config $name
: ${yarr_enable:=NO}

run_rc_command "$1"
