# www/nghttp2/files/openrc-nghttpx.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf to enable nghttpx:
# nghttpx_flags (str):		Set to "" by default.
#				Extra flags passed to start command.
#				Use --user=foo here if you wish. Make sure
#				fetch-ocsp-response can read the certificate files and
#				log-file directives can write their files.

name="nghttpx"
description=""

command="%%PREFIX%%/bin/nghttpx"
pidfile="/var/run/${name}.pid"

# read configuration and set defaults
: ${nghttpx_configfile:=%%PREFIX%%/etc/nghttpx/${name}.conf}
: ${nghttpx_env:="PATH=/bin:/usr/bin:%%PREFIX%%/bin"}	# Must include %%PREFIX%%/bin for python and fetch-ocsp-response

required_files="$nghttpx_configfile"

command_args="--conf ${nghttpx_configfile} --pid-file=${pidfile}"

depend() {
	keyword -shutdown
}
