# japanese/canna-lib/pkg-plist
bin/addwords
bin/canlisp
bin/cannacheck
bin/cannastat
bin/catdic
bin/cpdic
bin/crfreq
bin/crxdic
bin/crxgram
bin/cshost
bin/ctow
bin/delwords
bin/dicar
bin/dpbindic
bin/dpromdic
bin/dpxdic
bin/forcpp
bin/forsort
bin/kpdic
bin/lsdic
bin/mergeword
bin/mkbindic
bin/mkdic
bin/mkromdic
bin/mvdic
bin/rmdic
bin/splitword
bin/syncdic
bin/wtoc
etc/rc.d/canna.sh.sample
include/canna/RK.h
include/canna/jrkanji.h
include/canna/keydef.h
include/canna/mfdef.h
lib/canna/default.canna
lib/canna/dic/canna/bushu.cld
lib/canna/dic/canna/chimei.ctd
lib/canna/dic/canna/dics.dir
lib/canna/dic/canna/fuzokugo.cbd
lib/canna/dic/canna/fuzokugo.cld
lib/canna/dic/canna/hojomwd.ctd
lib/canna/dic/canna/hojoswd.ctd
lib/canna/dic/canna/iroha.cbd
lib/canna/dic/canna/iroha.cld
lib/canna/dic/canna/kanasmpl.ctd
lib/canna/dic/canna/katakana.ctd
lib/canna/dic/canna/keishiki.ctd
lib/canna/dic/canna/necgaiji.ctd
lib/canna/dic/canna/number.ctd
lib/canna/dic/canna/software.ctd
lib/canna/dic/canna/suffix.ctd
lib/canna/dic/default.cbp
lib/canna/dic/english.cbp
lib/canna/dic/jdaemon.cbp
lib/canna/dic/just.cbp
lib/canna/dic/kaisoku.cbp
lib/canna/dic/kana.cbp
lib/canna/dic/lan5.cbp
lib/canna/dic/matsu.cbp
lib/canna/dic/newjis.cbp
lib/canna/dic/romaji.cbp
lib/canna/dic/sokuon.cbp
lib/canna/dic/tut.cbp
lib/canna/dic/vje.cbp
lib/canna/dic/wx2+.cbp
lib/canna/sample/1.1.canna
lib/canna/sample/1.2.canna
lib/canna/sample/default.canna
lib/canna/sample/fixer.canna
lib/canna/sample/jdemon.canna
lib/canna/sample/just.canna
lib/canna/sample/lan5.canna
lib/canna/sample/matsu.canna
lib/canna/sample/modes.canna
lib/canna/sample/sample.canna
lib/canna/sample/skk.canna
lib/canna/sample/src/Makefile
lib/canna/sample/src/default.ctd
lib/canna/sample/src/english.ctd
lib/canna/sample/src/is.c
lib/canna/sample/src/jdaemon.ctd
lib/canna/sample/src/just.ctd
lib/canna/sample/src/kaisoku.ctd
lib/canna/sample/src/kana.ctd
lib/canna/sample/src/lan5.ctd
lib/canna/sample/src/matsu.ctd
lib/canna/sample/src/newjis.ctd
lib/canna/sample/src/romaji.ctd
lib/canna/sample/src/sokuon.ctd
lib/canna/sample/src/tut.ctd
lib/canna/sample/src/vje.ctd
lib/canna/sample/src/wx2+.ctd
lib/canna/sample/tut.canna
lib/canna/sample/unix.canna
lib/canna/sample/vje.canna
lib/canna/sample/wx2+.canna
lib/libRKC.a
lib/libRKC.so
lib/libRKC.so.1
lib/libRKC16.a
lib/libRKC16.so
lib/libRKC16.so.1
lib/libcanna.a
lib/libcanna.so
lib/libcanna.so.1
lib/libcanna16.a
lib/libcanna16.so
lib/libcanna16.so.1
sbin/cannaserver
sbin/cannakill
@exec mkdir -p %D/lib/canna/dic/user/user %D/lib/canna/dic/group
@exec chown bin:bin %D/lib/canna/dic/user/user %D/lib/canna/dic/user %D/lib/canna/dic/group
@exec chmod 775 %D/lib/canna/dic/user/user %D/lib/canna/dic/user %D/lib/canna/dic/group
@dirrm lib/canna/sample/src
@dirrm lib/canna/sample
@dirrm lib/canna/dic/canna
@dirrm include/canna
@unexec rmdir %D/lib/canna/dic/user/user 2>/dev/null || true
@unexec rmdir %D/lib/canna/dic/user      2>/dev/null || true
@unexec rmdir %D/lib/canna/dic/group     2>/dev/null || true
@unexec rmdir %D/lib/canna/dic           2>/dev/null || true
@unexec rmdir %D/lib/canna               2>/dev/null || true
