# net/mydns/pkg-plist
bin/mydnscheck
bin/mydnsexport
bin/mydnsimport
etc/mydns.conf.sample
etc/rc.d/mydns.sh.sample
@unexec install-info --delete %D/info/mydns.info %D/info/dir
info/mydns.info
@exec install-info %D/info/mydns.info %D/info/dir
sbin/mydns
share/locale/es/LC_MESSAGES/mydns.mo
share/locale/fr/LC_MESSAGES/mydns.mo
share/locale/pt_BR/LC_MESSAGES/mydns.mo
%%PORTDOCS%%share/doc/mydns/AUTHORS
%%PORTDOCS%%share/doc/mydns/COPYING
%%PORTDOCS%%share/doc/mydns/ChangeLog
%%PORTDOCS%%share/doc/mydns/NEWS
%%PORTDOCS%%share/doc/mydns/QUICKSTART
%%PORTDOCS%%share/doc/mydns/README
%%PORTDOCS%%share/doc/mydns/TODO
%%PORTDOCS%%share/doc/mydns/contrib/README
%%PORTDOCS%%share/doc/mydns/contrib/admin.php
%%PORTDOCS%%@dirrm share/doc/mydns/contrib
%%PORTDOCS%%@dirrm share/doc/mydns
