# net/haproxy-devel/pkg-plist
sbin/haproxy
etc/haproxy.cfg.sample
etc/rc.d/haproxy.sh.sample
