# net/haproxy-devel/pkg-plist
sbin/haproxy
etc/haproxy.cfg.sample
etc/rc.d/haproxy.sh.sample
%%PORTDOCS%%%%DOCSDIR%%/haproxy-en.txt
%%PORTDOCS%%%%DOCSDIR%%/haproxy-fr.txt
@dirrm %%DOCSDIR%%
