# sysutils/ssync/pkg-plist
bin/ssync
sbin/ssyncd
etc/ssyncd.conf.sample
etc/ssyncd.work.sample
etc/rc.d/ssyncd.init
share/man/man1/ssync.1.gz
