# audio/cadence/pkg-plist
bin/cadence
bin/cadence-aloop-daemon
bin/cadence-jackmeter
bin/cadence-jacksettings
bin/cadence-logs
bin/cadence-pulse2jack
bin/cadence-pulse2loopback
bin/cadence-render
bin/cadence-session-start
bin/cadence-xycontroller
bin/catarina
bin/catia
bin/claudia
bin/claudia-launcher
etc/X11/xinit/xinitrc.d/61cadence-session-inject
etc/xdg/autostart/cadence-session-start.desktop
share/applications/cadence.desktop
share/applications/catarina.desktop
share/applications/catia.desktop
share/applications/claudia-launcher.desktop
share/applications/claudia.desktop
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/arctican_plugins.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/ardour.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/aria.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/arpage.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/azr3.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_BME700.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_aks.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_arp2600.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_axxe.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_bassmaker.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_bit99.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_bitone.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_cs80.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_dx.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_explorer.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_granular.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_hammondB3.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_juno.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_jupiter8.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_memoryMoog.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_mini.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_mixer.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_monopoly.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_ms20.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_obx.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_obxa.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_odyssey.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_poly.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_poly800.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_pro1.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_prophet.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_prophet10.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_prophet52.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_realistic.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_rhodes.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_rhodesbass.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_roadrunner.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_sidney.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_solina.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_sonic6.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_stratus.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_trilogy.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_vox.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/bristol_voyager.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/distrho_plugins.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/energyxt2.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/foo-yc20.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/ingen.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/jack-rack.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/jacker.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/jsampler.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/juced_plugins.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/linuxdsp.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/lmms.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/loomer.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/lv2.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/mixxx.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/muse.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/non-mixer.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/non-sequencer.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/non-timeline.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/phasex.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/pianoteq.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/qamix.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/radium.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/rakarrack.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/renoise.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/tal_plugins.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/traverso.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/wolpertinger.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/yoshimi.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/zonage.png
%%DATADIR%%/icons/claudia-hicolor/16x16/apps/zynjacku.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/arctican_plugins.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/ardour.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/aria.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/arpage.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/azr3.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_BME700.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_aks.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_arp2600.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_axxe.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_bassmaker.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_bit100.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_bit99.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_bitone.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_cs80.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_dx.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_explorer.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_granular.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_hammondB3.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_juno.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_jupiter8.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_memoryMoog.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_mini.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_mixer.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_monopoly.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_ms20.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_obx.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_obxa.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_odyssey.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_poly.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_poly800.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_pro1.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_prophet.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_prophet10.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_prophet52.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_realistic.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_rhodes.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_rhodesbass.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_roadrunner.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_sidney.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_solina.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_sonic6.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_stratus.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_trilogy.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_vox.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_voxM2.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/bristol_voyager.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/distrho_plugins.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/energyxt2.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/foo-yc20.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/ingen.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/jack-rack.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/jacker.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/jsampler.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/juced_plugins.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/linuxdsp.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/lmms.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/loomer.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/lv2.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/mixxx.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/muse.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/non-mixer.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/non-sequencer.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/non-timeline.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/phasex.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/pianoteq.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/qamix.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/radium.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/rakarrack.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/renoise.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/tal_plugins.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/traverso.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/wolpertinger.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/yoshimi.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/zonage.png
%%DATADIR%%/icons/claudia-hicolor/48x48/apps/zynjacku.png
%%DATADIR%%/icons/claudia-hicolor/index.theme
%%DATADIR%%/pulse2jack/play+rec.pa
%%DATADIR%%/pulse2jack/play.pa
%%DATADIR%%/pulse2loopback/play+rec.pa
%%DATADIR%%/pulse2loopback/play.pa
%%DATADIR%%/src/cadence.py
%%DATADIR%%/src/cadence_aloop_daemon.py
%%DATADIR%%/src/cadence_session_start.py
%%DATADIR%%/src/canvaspreviewframe.py
%%DATADIR%%/src/catarina.py
%%DATADIR%%/src/catia.py
%%DATADIR%%/src/claudia.py
%%DATADIR%%/src/claudia_database.py
%%DATADIR%%/src/claudia_launcher.py
%%DATADIR%%/src/clickablelabel.py
%%DATADIR%%/src/jacklib.py
%%DATADIR%%/src/jacklib_helpers.py
%%DATADIR%%/src/jacksettings.py
%%DATADIR%%/src/logs.py
%%DATADIR%%/src/patchcanvas.py
%%DATADIR%%/src/patchcanvas_theme.py
%%DATADIR%%/src/render.py
%%DATADIR%%/src/resources_rc.py
%%DATADIR%%/src/shared.py
%%DATADIR%%/src/shared_cadence.py
%%DATADIR%%/src/shared_canvasjack.py
%%DATADIR%%/src/shared_settings.py
%%DATADIR%%/src/systray.py
%%DATADIR%%/src/ui_cadence.py
%%DATADIR%%/src/ui_cadence_rwait.py
%%DATADIR%%/src/ui_cadence_tb_a2j.py
%%DATADIR%%/src/ui_cadence_tb_alsa.py
%%DATADIR%%/src/ui_cadence_tb_jack.py
%%DATADIR%%/src/ui_cadence_tb_pa.py
%%DATADIR%%/src/ui_catarina.py
%%DATADIR%%/src/ui_catarina_addgroup.py
%%DATADIR%%/src/ui_catarina_addport.py
%%DATADIR%%/src/ui_catarina_connectports.py
%%DATADIR%%/src/ui_catarina_disconnectports.py
%%DATADIR%%/src/ui_catarina_removegroup.py
%%DATADIR%%/src/ui_catarina_removeport.py
%%DATADIR%%/src/ui_catarina_renamegroup.py
%%DATADIR%%/src/ui_catarina_renameport.py
%%DATADIR%%/src/ui_catia.py
%%DATADIR%%/src/ui_claudia.py
%%DATADIR%%/src/ui_claudia_createroom.py
%%DATADIR%%/src/ui_claudia_launcher.py
%%DATADIR%%/src/ui_claudia_launcher_app.py
%%DATADIR%%/src/ui_claudia_projectname.py
%%DATADIR%%/src/ui_claudia_projectproperties.py
%%DATADIR%%/src/ui_claudia_runcustom.py
%%DATADIR%%/src/ui_claudia_studiolist.py
%%DATADIR%%/src/ui_claudia_studioname.py
%%DATADIR%%/src/ui_logs.py
%%DATADIR%%/src/ui_render.py
%%DATADIR%%/src/ui_settings_app.py
%%DATADIR%%/src/ui_settings_jack.py
%%DATADIR%%/templates/Ardour2/Ardour2.ardour
%%DATADIR%%/templates/Ardour2/instant.xml
%%DATADIR%%/templates/Ardour3/Ardour3.ardour
%%DATADIR%%/templates/Ardour3/instant.xml
%%DATADIR%%/templates/CalfJackHost
%%DATADIR%%/templates/Carla.carxp
%%DATADIR%%/templates/Hydrogen.h2song
%%DATADIR%%/templates/Jack-Mixer.xml
%%DATADIR%%/templates/Jack-Rack.xml
%%DATADIR%%/templates/Jacker.jsong
%%DATADIR%%/templates/LMMS.mmp
%%DATADIR%%/templates/MusE.med
%%DATADIR%%/templates/Non-Sequencer.non
%%DATADIR%%/templates/Non-Timeline/history
%%DATADIR%%/templates/Non-Timeline/info
%%DATADIR%%/templates/Qsampler.lscp
%%DATADIR%%/templates/Qtractor.qtr
%%DATADIR%%/templates/Rosegarden.rg
%%DATADIR%%/templates/Seq24.midi
share/icons/hicolor/128x128/apps/cadence.png
share/icons/hicolor/128x128/apps/catarina.png
share/icons/hicolor/128x128/apps/catia.png
share/icons/hicolor/128x128/apps/claudia-launcher.png
share/icons/hicolor/128x128/apps/claudia.png
share/icons/hicolor/16x16/apps/cadence.png
share/icons/hicolor/16x16/apps/catarina.png
share/icons/hicolor/16x16/apps/catia.png
share/icons/hicolor/16x16/apps/claudia-launcher.png
share/icons/hicolor/16x16/apps/claudia.png
share/icons/hicolor/256x256/apps/cadence.png
share/icons/hicolor/256x256/apps/catarina.png
share/icons/hicolor/256x256/apps/catia.png
share/icons/hicolor/256x256/apps/claudia-launcher.png
share/icons/hicolor/256x256/apps/claudia.png
share/icons/hicolor/48x48/apps/cadence.png
share/icons/hicolor/48x48/apps/catarina.png
share/icons/hicolor/48x48/apps/catia.png
share/icons/hicolor/48x48/apps/claudia-launcher.png
share/icons/hicolor/48x48/apps/claudia.png
share/icons/hicolor/scalable/apps/cadence.svg
share/icons/hicolor/scalable/apps/catarina.svg
share/icons/hicolor/scalable/apps/catia.svg
share/icons/hicolor/scalable/apps/claudia-launcher.svg
share/icons/hicolor/scalable/apps/claudia.svg
