# dns/mydns-ng/pkg-plist
etc/mydns.conf.sample
etc/rc.d/mydns.sh.sample
@unexec install-info --delete %D/info/mydns.info %D/info/dir
info/mydns.info
@exec install-info %D/info/mydns.info %D/info/dir
bin/mydnscheck
bin/mydnsexport
bin/mydnsimport
share/locale/es/LC_MESSAGES/mydns.mo
share/locale/pt_BR/LC_MESSAGES/mydns.mo
share/locale/fr/LC_MESSAGES/mydns.mo
sbin/mydns
%%PORTDOCS%%share/doc/mydns/AUTHORS
%%PORTDOCS%%share/doc/mydns/COPYING
%%PORTDOCS%%share/doc/mydns/ChangeLog
%%PORTDOCS%%share/doc/mydns/NEWS
%%PORTDOCS%%share/doc/mydns/QUICKSTART
%%PORTDOCS%%share/doc/mydns/README
%%PORTDOCS%%share/doc/mydns/TODO
%%PORTDOCS%%share/doc/mydns/contrib/README
%%PORTDOCS%%share/doc/mydns/contrib/admin.php
%%PORTDOCS%%share/doc/mydns/html/mydns.html
%%PORTDOCS%%share/doc/mydns/html/mydns_1.html
%%PORTDOCS%%share/doc/mydns/html/mydns_10.html
%%PORTDOCS%%share/doc/mydns/html/mydns_11.html
%%PORTDOCS%%share/doc/mydns/html/mydns_12.html
%%PORTDOCS%%share/doc/mydns/html/mydns_13.html
%%PORTDOCS%%share/doc/mydns/html/mydns_14.html
%%PORTDOCS%%share/doc/mydns/html/mydns_15.html
%%PORTDOCS%%share/doc/mydns/html/mydns_16.html
%%PORTDOCS%%share/doc/mydns/html/mydns_17.html
%%PORTDOCS%%share/doc/mydns/html/mydns_18.html
%%PORTDOCS%%share/doc/mydns/html/mydns_19.html
%%PORTDOCS%%share/doc/mydns/html/mydns_2.html
%%PORTDOCS%%share/doc/mydns/html/mydns_20.html
%%PORTDOCS%%share/doc/mydns/html/mydns_21.html
%%PORTDOCS%%share/doc/mydns/html/mydns_22.html
%%PORTDOCS%%share/doc/mydns/html/mydns_23.html
%%PORTDOCS%%share/doc/mydns/html/mydns_24.html
%%PORTDOCS%%share/doc/mydns/html/mydns_25.html
%%PORTDOCS%%share/doc/mydns/html/mydns_26.html
%%PORTDOCS%%share/doc/mydns/html/mydns_27.html
%%PORTDOCS%%share/doc/mydns/html/mydns_28.html
%%PORTDOCS%%share/doc/mydns/html/mydns_29.html
%%PORTDOCS%%share/doc/mydns/html/mydns_3.html
%%PORTDOCS%%share/doc/mydns/html/mydns_30.html
%%PORTDOCS%%share/doc/mydns/html/mydns_31.html
%%PORTDOCS%%share/doc/mydns/html/mydns_32.html
%%PORTDOCS%%share/doc/mydns/html/mydns_33.html
%%PORTDOCS%%share/doc/mydns/html/mydns_34.html
%%PORTDOCS%%share/doc/mydns/html/mydns_4.html
%%PORTDOCS%%share/doc/mydns/html/mydns_5.html
%%PORTDOCS%%share/doc/mydns/html/mydns_6.html
%%PORTDOCS%%share/doc/mydns/html/mydns_7.html
%%PORTDOCS%%share/doc/mydns/html/mydns_8.html
%%PORTDOCS%%share/doc/mydns/html/mydns_9.html
%%PORTDOCS%%share/doc/mydns/html/mydns_abt.html
%%PORTDOCS%%share/doc/mydns/html/mydns_ovr.html
%%PORTDOCS%%share/doc/mydns/html/mydns_toc.html
%%PORTDOCS%%@dirrm share/doc/mydns/html
%%PORTDOCS%%@dirrm share/doc/mydns/contrib
%%PORTDOCS%%@dirrm share/doc/mydns
