# net/trafcount/pkg-plist
bin/trafcount
etc/rc.d/trafcount.sh
libexec/trafcount.ko
