# security/fwknop/files/openrc-fwknopd.in
#!/sbin/openrc-run

name=fwknopd

command=%%PREFIX%%/sbin/fwknopd
required_files=%%PREFIX%%/etc/fwknop/access.conf
