# audio/firefly/files/mt-daapd.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: firefly
# REQUIRE: LOGIN

. /etc/rc.subr

: ${firefly_enable=NO}

name="firefly"
rcvar="`set_rcvar`"

command="%%PREFIX%%/sbin/mt-daapd"
command_args="-c %%PREFIX%%/etc/mt-daapd.conf"
pidfile="/var/run/mt-daapd.pid"
required_dirs="%%PREFIX%%/var/db/firefly"
required_files="%%PREFIX%%/etc/mt-daapd.conf"

load_rc_config "$name"
run_rc_command "$1"
