# audio/firefly/files/openrc-mt-daapd.in
#!/sbin/openrc-run

name="firefly"
command="%%PREFIX%%/sbin/mt-daapd"
command_args="-c %%PREFIX%%/etc/firefly/mt-daapd.conf"
pidfile="/var/run/mt-daapd.pid"
required_dirs="%%DAAPD_DBDIR%%"
required_files="%%PREFIX%%/etc/firefly/mt-daapd.conf"
