# net-mgmt/pandorafms_server/pkg-plist
%%PERL5_MAN3%%/PandoraFMS::Core.3.gz
%%PERL5_MAN3%%/PandoraFMS::GIS.3.gz
%%PERL5_MAN3%%/PandoraFMS::GeoIP.3.gz
%%PERL5_MAN3%%/PandoraFMS::NmapParser.3.gz
%%PERL5_MAN3%%/PandoraFMS::Sendmail.3.gz
%%PERL5_MAN3%%/PandoraFMS::Tools.3.gz
%%PERL5_MAN3%%/PandoraFMS::Traceroute.3.gz
%%PERL5_MAN3%%/PandoraFMS::Traceroute::PurePerl.3.gz
%%SITE_PERL%%/PandoraFMS/Config.pm
%%SITE_PERL%%/PandoraFMS/Core.pm
%%SITE_PERL%%/PandoraFMS/DB.pm
%%SITE_PERL%%/PandoraFMS/DataServer.pm
%%SITE_PERL%%/PandoraFMS/GIS.pm
%%SITE_PERL%%/PandoraFMS/GeoIP.pm
%%SITE_PERL%%/PandoraFMS/NetworkServer.pm
%%SITE_PERL%%/PandoraFMS/NmapParser.pm
%%SITE_PERL%%/PandoraFMS/PluginServer.pm
%%SITE_PERL%%/PandoraFMS/PredictionServer.pm
%%SITE_PERL%%/PandoraFMS/ProducerConsumerServer.pm
%%SITE_PERL%%/PandoraFMS/ReconServer.pm
%%SITE_PERL%%/PandoraFMS/SNMPServer.pm
%%SITE_PERL%%/PandoraFMS/Sendmail.pm
%%SITE_PERL%%/PandoraFMS/Server.pm
%%SITE_PERL%%/PandoraFMS/Tools.pm
%%SITE_PERL%%/PandoraFMS/Traceroute.pm
%%SITE_PERL%%/PandoraFMS/Traceroute/PurePerl.pm
%%SITE_PERL%%/PandoraFMS/WMIServer.pm
man/man1/pandora_server.1.gz
man/man1/tentacle_server.1.gz
%%DATADIR%%/util/alert_multicast.pl
%%DATADIR%%/util/change_remoteconfig.pl
%%DATADIR%%/util/compaq_chassis_trap_manager.pl
%%DATADIR%%/util/compaq_chassis_trap_manager2.pl
%%DATADIR%%/util/compare_schemas.pl
%%DATADIR%%/util/esx_trap_manager.pl
%%DATADIR%%/util/gis.README
%%DATADIR%%/util/gpx2pandora_agent_data.pl
%%DATADIR%%/util/integria_rticket.pl
%%DATADIR%%/util/integrity_check
%%DATADIR%%/util/n2p.README
%%DATADIR%%/util/n2p.pl
%%DATADIR%%/util/pandora2ast/GPL
%%DATADIR%%/util/pandora2ast/README
%%DATADIR%%/util/pandora2ast/channelcommands
%%DATADIR%%/util/pandora2ast/extensions.conf.add
%%DATADIR%%/util/pandora2ast/pandora2ast.pl
%%DATADIR%%/util/pandora2ast/text2speech.php
%%DATADIR%%/util/pandora_DBI_test.pl
%%DATADIR%%/util/pandora_backup.sh
%%DATADIR%%/util/pandora_count.sh
%%DATADIR%%/util/pandora_csvbridge.pl
%%DATADIR%%/util/pandora_database_check.pl
%%DATADIR%%/util/pandora_db.pl
%%DATADIR%%/util/pandora_dbstress.README
%%DATADIR%%/util/pandora_dbstress.pl
%%DATADIR%%/util/pandora_diagnostic.sh
%%DATADIR%%/util/pandora_logrotate
%%DATADIR%%/util/pandora_manage.pl
%%DATADIR%%/util/pandora_migrate_plugins.pl
%%DATADIR%%/util/pandora_migrate_recon_scripts.pl
%%DATADIR%%/util/pandora_recode_db.pl
%%DATADIR%%/util/pandora_remote_agent.sh
%%DATADIR%%/util/pandora_revent.pl
%%DATADIR%%/util/pandora_revent_create.pl
%%DATADIR%%/util/pandora_server
%%DATADIR%%/util/pandora_server_logrotate
%%DATADIR%%/util/pandora_snmp.README
%%DATADIR%%/util/pandora_snmp.sh
%%DATADIR%%/util/pandora_sync_db.pl
%%DATADIR%%/util/pandora_watchdog.sh
%%DATADIR%%/util/pandora_xml_stress.README
%%DATADIR%%/util/pandora_xml_stress.agents
%%DATADIR%%/util/pandora_xml_stress.conf
%%DATADIR%%/util/pandora_xml_stress.pl
%%DATADIR%%/util/pandora_xml_stress_module_source.txt
%%DATADIR%%/util/plugin/EU_10yrspread.pl
%%DATADIR%%/util/plugin/SMTP_check.pl
%%DATADIR%%/util/plugin/babel_plugin/README
%%DATADIR%%/util/plugin/babel_plugin/babel_plugin.sh
%%DATADIR%%/util/plugin/babel_plugin/plugin_definition.ini
%%DATADIR%%/util/plugin/bgp_plugin.pl
%%DATADIR%%/util/plugin/create_integria_incident.sh
%%DATADIR%%/util/plugin/dns_plugin.sh
%%DATADIR%%/util/plugin/dummy_plugin.pl
%%DATADIR%%/util/plugin/iface_bandwith.pl
%%DATADIR%%/util/plugin/integria_plugin/README
%%DATADIR%%/util/plugin/integria_plugin/integria_plugin.sh
%%DATADIR%%/util/plugin/integria_plugin/plugin_definition.ini
%%DATADIR%%/util/plugin/ipmi-plugin.pl
%%DATADIR%%/util/plugin/multicast.pl
%%DATADIR%%/util/plugin/mysql_plugin.sh
%%DATADIR%%/util/plugin/openvpn_pandoraplugin.pl
%%DATADIR%%/util/plugin/packet_loss.sh
%%DATADIR%%/util/plugin/pandora_inventory_change.README
%%DATADIR%%/util/plugin/pandora_inventory_change.pl
%%DATADIR%%/util/plugin/pandora_loadgen.pl
%%DATADIR%%/util/plugin/pandora_server_status.pl
%%DATADIR%%/util/plugin/snmp_process.pl
%%DATADIR%%/util/plugin/snmp_remote.pl
%%DATADIR%%/util/plugin/ssh_pandoraplugin.sh
%%DATADIR%%/util/plugin/udp_nmap_plugin.sh
%%DATADIR%%/util/plugin/webcheck_plugin.sh
%%DATADIR%%/util/recon_scripts/ipmi-recon.pl
%%DATADIR%%/util/recon_scripts/snmp-recon.pl
%%DATADIR%%/util/recon_scripts/wmi-recon.pl
%%DATADIR%%/util/snmptrap_gen.sh
%%DATADIR%%/util/tentacle_serverd
%%DATADIR%%/util/trap_rate.sh
%%DATADIR%%/util/udp_client.pl
%%DATADIR%%/util/windows.README
@mode 555
bin/pandora_exec
bin/pandora_server
%%TENTACLE%%bin/tentacle_server
etc/rc.d/pandora_server
etc/rc.d/tentacle_server
@sample(,,440) %%ETCDIR%%/pandora_server.conf.sample
