# net/openistgt/pkg-plist
@comment $FreeBSD$
bin/istgt
bin/istgtcontrol
etc/rc.d/istgt
etc/istgt/istgt.conf.sample
etc/istgt/auth.conf.sample
etc/istgt/istgtcontrol.conf.sample
@dirrm etc/istgt
