# databases/memcached/pkg-plist
bin/memcached
etc/rc.d/memcached.sh
%%PORTDOCS%%%%DOCSDIR%%/memory_management.txt
%%PORTDOCS%%%%DOCSDIR%%/protocol.txt
%%PORTDOCS%%@dirrm %%DOCSDIR%%
