# devel/p4ftpd/pkg-plist
bin/p4
bin/p4web
sbin/p4d
etc/perforce
etc/rc.d/perforce.sh
