# emulators/vmware3/pkg-plist
bin/vmware
etc/rc.d/vmware.sh
etc/vmware/config
lib/vmware/lib/modules/vmmon_smp.ko
%%FREEBSD5%%lib/vmware/lib/modules/linker.hints
lib/vmware/lib/modules/vmmon_up.ko
lib/vmware/lib/modules/vmnet.ko
lib/vmware/lib/config
lib/vmware/lib/configurator/vmnet-dhcpd.conf
lib/vmware/lib/configurator/vmnet-nat.conf
lib/vmware/lib/configurator/vmnet-smb.conf
lib/vmware/lib/isoimages/freebsd.iso
lib/vmware/lib/isoimages/linux.iso
lib/vmware/lib/isoimages/windows.iso
lib/vmware/lib/xkeymap/be101
lib/vmware/lib/xkeymap/be104
lib/vmware/lib/xkeymap/chde101
lib/vmware/lib/xkeymap/chde104
lib/vmware/lib/xkeymap/chfr101
lib/vmware/lib/xkeymap/chfr104
lib/vmware/lib/xkeymap/de101
lib/vmware/lib/xkeymap/de104
lib/vmware/lib/xkeymap/dk101
lib/vmware/lib/xkeymap/dk104
lib/vmware/lib/xkeymap/es101
lib/vmware/lib/xkeymap/es104
lib/vmware/lib/xkeymap/fi101
lib/vmware/lib/xkeymap/fi104
lib/vmware/lib/xkeymap/fr101
lib/vmware/lib/xkeymap/fr104
lib/vmware/lib/xkeymap/gb101
lib/vmware/lib/xkeymap/gb104
lib/vmware/lib/xkeymap/is101
lib/vmware/lib/xkeymap/is104
lib/vmware/lib/xkeymap/it101
lib/vmware/lib/xkeymap/it104
lib/vmware/lib/xkeymap/jp106
lib/vmware/lib/xkeymap/jp109
lib/vmware/lib/xkeymap/no101
lib/vmware/lib/xkeymap/no104
lib/vmware/lib/xkeymap/pt101
lib/vmware/lib/xkeymap/pt104
lib/vmware/lib/xkeymap/se101
lib/vmware/lib/xkeymap/se104
lib/vmware/lib/xkeymap/us101
lib/vmware/lib/xkeymap/us104
lib/vmware/lib/help/PowerOff.gif
lib/vmware/lib/help/PowerOn.gif
lib/vmware/lib/help/PowerReset.gif
lib/vmware/lib/help/Resume.gif
lib/vmware/lib/help/Suspend.gif
lib/vmware/lib/help/Wizcd.htm
lib/vmware/lib/help/Wizcongrats.htm
lib/vmware/lib/help/Wizdisk.htm
lib/vmware/lib/help/Wizdisk3.htm
lib/vmware/lib/help/WizdiskType.htm
lib/vmware/lib/help/Wizfloppy.htm
lib/vmware/lib/help/Wizguestkit.htm
lib/vmware/lib/help/Wiznetwork.htm
lib/vmware/lib/help/Wizos.htm
lib/vmware/lib/help/Wizpath.htm
lib/vmware/lib/help/about_this_help.htm
lib/vmware/lib/help/commandref_linux.htm
lib/vmware/lib/help/config_editor.htm
lib/vmware/lib/help/config_editor_newvm.htm
lib/vmware/lib/help/configvm.htm
lib/vmware/lib/help/configvm_adv.htm
lib/vmware/lib/help/configvm_automation.htm
lib/vmware/lib/help/configvm_debug.htm
lib/vmware/lib/help/configvm_disableaccel.htm
lib/vmware/lib/help/configvm_display-problems.htm
lib/vmware/lib/help/configvm_displayname.htm
lib/vmware/lib/help/configvm_guestos.htm
lib/vmware/lib/help/configvm_logs.htm
lib/vmware/lib/help/configvm_memory.htm
lib/vmware/lib/help/configvm_options.htm
lib/vmware/lib/help/configvm_state.htm
lib/vmware/lib/help/configvm_virtualdisk.htm
lib/vmware/lib/help/configvmware-disablemem.htm
lib/vmware/lib/help/configvmware-inputgrab.htm
lib/vmware/lib/help/configvmware-maxmemory.htm
lib/vmware/lib/help/configvmware-shift.htm
lib/vmware/lib/help/configvmware.htm
lib/vmware/lib/help/configvmware_copy_paste.htm
lib/vmware/lib/help/devices.htm
lib/vmware/lib/help/devices_autostart.htm
lib/vmware/lib/help/devices_cdrom-changeconnect.htm
lib/vmware/lib/help/devices_cdrom-disable.htm
lib/vmware/lib/help/devices_cdrom-issues-linux.htm
lib/vmware/lib/help/devices_cdrom-issues.htm
lib/vmware/lib/help/devices_cdrom-removable.htm
lib/vmware/lib/help/devices_cdrom-remove.htm
lib/vmware/lib/help/devices_cdrom.htm
lib/vmware/lib/help/devices_floppy.htm
lib/vmware/lib/help/devices_floppy_changeconnect.htm
lib/vmware/lib/help/devices_genscsi.htm
lib/vmware/lib/help/devices_genscsi_linux.htm
lib/vmware/lib/help/devices_modem.htm
lib/vmware/lib/help/devices_mouse.htm
lib/vmware/lib/help/devices_netadapter.htm
lib/vmware/lib/help/devices_parallel-info.htm
lib/vmware/lib/help/devices_parallel-onguestos.htm
lib/vmware/lib/help/devices_parallel.htm
lib/vmware/lib/help/devices_parallel_changeconnect.htm
lib/vmware/lib/help/devices_pcmcia.htm
lib/vmware/lib/help/devices_plaindisk-info.htm
lib/vmware/lib/help/devices_plaindisk.htm
lib/vmware/lib/help/devices_printer.htm
lib/vmware/lib/help/devices_rawdisk.htm
lib/vmware/lib/help/devices_remove_netadapter.htm
lib/vmware/lib/help/devices_scanner.htm
lib/vmware/lib/help/devices_scsi-bios.htm
lib/vmware/lib/help/devices_scsi-guestos.htm
lib/vmware/lib/help/devices_scsi.htm
lib/vmware/lib/help/devices_scsi_info.htm
lib/vmware/lib/help/devices_scsi_vs-ide.htm
lib/vmware/lib/help/devices_serial.htm
lib/vmware/lib/help/devices_sound-issues.htm
lib/vmware/lib/help/devices_sound.htm
lib/vmware/lib/help/devices_sound_changeconn.htm
lib/vmware/lib/help/devices_usb-info.htm
lib/vmware/lib/help/devices_usb.htm
lib/vmware/lib/help/devices_video-linux-res.htm
lib/vmware/lib/help/devices_video-linux.htm
lib/vmware/lib/help/glossary.htm
lib/vmware/lib/help/devices_video-win2000.htm
lib/vmware/lib/help/devices_video-win95.htm
lib/vmware/lib/help/devices_video-win98.htm
lib/vmware/lib/help/devices_virtualdrive-consolidate.htm
lib/vmware/lib/help/devices_virtualdrive-diskmode.htm
lib/vmware/lib/help/devices_virtualdrive-virtualdisk-physical.htm
lib/vmware/lib/help/devices_virtualdrive-virtualdisk-reuse.htm
lib/vmware/lib/help/devices_virtualdrive-virtualdisk.htm
lib/vmware/lib/help/devices_virtualdrive.htm
lib/vmware/lib/help/devices_zipdrive.htm
lib/vmware/lib/help/grab.htm
lib/vmware/lib/help/guestos-corellinux.htm
lib/vmware/lib/help/guestos-dos.htm
lib/vmware/lib/help/guestos.htm
lib/vmware/lib/help/guestos-freebsd.htm
lib/vmware/lib/help/guestos-kit.htm
lib/vmware/lib/help/guestos-mandrake-linux.htm
lib/vmware/lib/help/guestos-mandrake-linux8081.htm
lib/vmware/lib/help/guestos-mandrake-linux82.htm
lib/vmware/lib/help/guestos-novell-4.htm
lib/vmware/lib/help/guestos-novell-5.htm
lib/vmware/lib/help/guestos-novell-6.htm
lib/vmware/lib/help/guestos-novell-zenwrks.htm
lib/vmware/lib/help/guestos-openlinux.htm
lib/vmware/lib/help/guestos-rawdisk.htm
lib/vmware/lib/help/guestos-redhat-linux.htm
lib/vmware/lib/help/guestos-redhat-linux6.x.htm
lib/vmware/lib/help/guestos-redhat-linux7.0.htm
lib/vmware/lib/help/guestos-redhat-linux7.1.htm
lib/vmware/lib/help/guestos-redhat-linux7.3.htm
lib/vmware/lib/help/guestos-solaris7.htm
lib/vmware/lib/help/guestos-suse-linux-6.x.htm
lib/vmware/lib/help/guestos-suse-linux-7.0.htm
lib/vmware/lib/help/guestos-suse-linux-7.1.htm
lib/vmware/lib/help/guestos-suse-linux-8.0.htm
lib/vmware/lib/help/guestos-suse-linux.htm
lib/vmware/lib/help/guestos-turbolinux-6.0.htm
lib/vmware/lib/help/guestos-turbolinux-7.0.htm
lib/vmware/lib/help/guestos-turbolinux.htm
lib/vmware/lib/help/guestos-win2000.htm
lib/vmware/lib/help/guestos-win95.htm
lib/vmware/lib/help/guestos-win98.htm
lib/vmware/lib/help/guestos-winNETSvr.htm
lib/vmware/lib/help/guestos-winXP.htm
lib/vmware/lib/help/hotkey.htm
lib/vmware/lib/help/guestos-winme.htm
lib/vmware/lib/help/guestos-winnt.htm
lib/vmware/lib/help/image001.gif
lib/vmware/lib/help/index.htm
lib/vmware/lib/help/install_apps.htm
lib/vmware/lib/help/introduction.htm
lib/vmware/lib/help/keyboard_issus.htm
lib/vmware/lib/help/lock_files.htm
lib/vmware/lib/help/memory.htm
lib/vmware/lib/help/migratevm-toupgrade.htm
lib/vmware/lib/help/migratevm-uninstall.htm
lib/vmware/lib/help/migratevm-upgrade.htm
lib/vmware/lib/help/migratevm.htm
lib/vmware/lib/help/migratevm_buildno.htm
lib/vmware/lib/help/migratevm_pc.htm
lib/vmware/lib/help/migrating.htm
lib/vmware/lib/help/movefiles.htm
lib/vmware/lib/help/networking-components.htm
lib/vmware/lib/help/newvm.htm
lib/vmware/lib/help/networking-onguestos.htm
lib/vmware/lib/help/networking-onhost.htm
lib/vmware/lib/help/networking-promiscuous.htm
lib/vmware/lib/help/networking.htm
lib/vmware/lib/help/networking_assignip.htm
lib/vmware/lib/help/networking_bridged-onguestos.htm
lib/vmware/lib/help/networking_bridged-second.htm
lib/vmware/lib/help/networking_bridged.htm
lib/vmware/lib/help/networking_changing.htm
lib/vmware/lib/help/networking_custom.htm
lib/vmware/lib/help/networking_host-onguestos.htm
lib/vmware/lib/help/networking_host-twointerfaces.htm
lib/vmware/lib/help/networking_host.htm
lib/vmware/lib/help/networking_host_ipleaks.htm
lib/vmware/lib/help/perform.htm
lib/vmware/lib/help/networking_host_issues-linux.htm
lib/vmware/lib/help/networking_host_routing.htm
lib/vmware/lib/help/networking_host_samba.htm
lib/vmware/lib/help/networking_ipmasqerade.htm
lib/vmware/lib/help/networking_ipmasqerade_firewallrh.htm
lib/vmware/lib/help/networking_ipmasqerade_fwconfig.htm
lib/vmware/lib/help/networking_macaddr.htm
lib/vmware/lib/help/networking_nat.htm
lib/vmware/lib/help/networking_nat_intro.htm
lib/vmware/lib/help/networking_proxy.htm
lib/vmware/lib/help/networking_wireless.htm
lib/vmware/lib/help/newvm-diskpartition.htm
lib/vmware/lib/help/newvm-diskpartition_access.htm
lib/vmware/lib/help/newvm-dualboot-hdwr.htm
lib/vmware/lib/help/newvm-diskpartition_hide.htm
lib/vmware/lib/help/newvm-dualboot-scsi.htm
lib/vmware/lib/help/newvm-dualboot.htm
lib/vmware/lib/help/newvm-virtualdisk.htm
lib/vmware/lib/help/newvm_wizard.htm
lib/vmware/lib/help/newvm_wizard_process_cdrom.htm
lib/vmware/lib/help/newvm_wizard_process_directory.htm
lib/vmware/lib/help/newvm_wizard_process_disksize.htm
lib/vmware/lib/help/newvm_wizard_process_disktype.htm
lib/vmware/lib/help/newvm_wizard_process_floppy.htm
lib/vmware/lib/help/newvm_wizard_process_guestos.htm
lib/vmware/lib/help/newvm_wizard_process_network.htm
lib/vmware/lib/help/perform-guestos-linux.htm
lib/vmware/lib/help/toc.htm
lib/vmware/lib/help/perform-guestos-win2000.htm
lib/vmware/lib/help/perform-guestos-win95-98.htm
lib/vmware/lib/help/perform-linuxhost.htm
lib/vmware/lib/help/screensavers.htm
lib/vmware/lib/help/serialnumber.htm
lib/vmware/lib/help/sharefiles_hos_gos.htm
lib/vmware/lib/help/sharevm.htm
lib/vmware/lib/help/sizewin.htm
lib/vmware/lib/help/sizewin_fullscreen.htm
lib/vmware/lib/help/startvm.htm
lib/vmware/lib/help/startvmware.htm
lib/vmware/lib/help/statevm.htm
lib/vmware/lib/help/statevm_off.htm
lib/vmware/lib/help/statevm_on.htm
lib/vmware/lib/help/statevm_repeat_resume.htm
lib/vmware/lib/help/statevm_reset.htm
lib/vmware/lib/help/statevm_resume.htm
lib/vmware/lib/help/statevm_suspend.htm
lib/vmware/lib/help/support.htm
lib/vmware/lib/help/switching.htm
lib/vmware/lib/help/sysreqs.htm
lib/vmware/lib/help/user-permissions-linux.htm
lib/vmware/lib/help/vm_common.htm
lib/vmware/lib/help/vmtools-copypaste.htm
lib/vmware/lib/help/vmtools-defrag.htm
lib/vmware/lib/help/vmtools-shrink.htm
lib/vmware/lib/help/vmtools.htm
lib/vmware/lib/help/vmtools_install.htm
lib/vmware/lib/help/vmtools_install_freebsd.htm
lib/vmware/lib/help/vmtools_install_linux.htm
lib/vmware/lib/help/vmtools_install_win.htm
lib/vmware/lib/help/vmtools_prefs.htm
lib/vmware/lib/help/whatsnew.htm
lib/vmware/lib/help/CUSTOM-NETWORK.gif
lib/vmware/lib/help/ConfigWiz.htm
lib/vmware/lib/help/FullScreen.gif
lib/vmware/lib/help/Help.gif
lib/vmware/bin/df
lib/vmware/lib/messages/chef/vmware-mks.vmsg
lib/vmware/lib/messages/chef/vmware-ui.vmsg
lib/vmware/lib/messages/chef/vmware.vmsg
lib/vmware/lib/messages/ja/vmware-mks.vmsg
lib/vmware/lib/messages/ja/vmware-ui.vmsg
lib/vmware/lib/messages/ja/vmware.vmsg
lib/vmware/lib/messages/jive/vmware-mks.vmsg
lib/vmware/lib/messages/jive/vmware-ui.vmsg
lib/vmware/lib/messages/jive/vmware.vmsg
lib/vmware/bin/vmware-config.pl
lib/vmware/bin/vmware-uninstall.pl
lib/vmware/bin/vmnet-bridge
lib/vmware/bin/vmnet-dhcpd
lib/vmware/bin/vmnet-sniffer
lib/vmware/bin/vmware-ping
lib/vmware/bin/vmware-ui
lib/vmware/bin/vmware-mks
lib/vmware/bin/vmware-wizard
lib/vmware/bin/vmware
lib/vmware/bin/vmnet-natd
lib/vmware/bin/vmnet-netifup
lib/vmware/bin/vmware-nmbd
lib/vmware/bin/vmware-smbd
lib/vmware/bin/vmware-smbpasswd
lib/vmware/bin/vmware-smbpasswd.bin
lib/vmware/lib/licenses/user/license.ws.3.0
lib/vmware/lib/smb/codepages/codepage.437
lib/vmware/lib/smb/codepages/codepage.737
lib/vmware/lib/smb/codepages/codepage.850
lib/vmware/lib/smb/codepages/codepage.852
lib/vmware/lib/smb/codepages/codepage.861
lib/vmware/lib/smb/codepages/codepage.866
lib/vmware/lib/smb/codepages/codepage.932
lib/vmware/lib/smb/codepages/codepage.936
lib/vmware/lib/smb/codepages/codepage.949
lib/vmware/lib/smb/codepages/codepage.950
share/doc/vmware/Hints.FreeBSD
share/doc/vmware/README.FreeBSD
share/doc/vmware/open_source_licenses.txt
share/doc/vmware/EULA
share/doc/vmware/INSTALL
share/doc/vmware/README
@exec mkdir -p %%LINUXBASE%%/dev;for n in 0 1 2 3 4 5 6 7 8 9; do ln -s /dev/ttyv$n %%LINUXBASE%%/dev/tty$n;done
@exec ln -s /dev/ttyva %%LINUXBASE%%/dev/tty10
@exec ln -s /dev/ttyvb %%LINUXBASE%%/dev/tty11
@exec mknod %%LINUXBASE%%/dev/null c 2 2
@exec chmod 666 %%LINUXBASE%%/dev/null
@exec mknod %%LINUXBASE%%/dev/hda b 0 0x00010002
@exec mknod %%LINUXBASE%%/dev/hdb b 0 0x0001000a
@exec ln -sf %D/etc/vmware %%LINUXBASE%%/etc/
@exec mknod %%LINUXBASE%%/dev/vmmon c 200 0
@exec mknod %%LINUXBASE%%/dev/vmnet1 c 149 %%VMNET1_MINOR%%
@unexec rm -f %%LINUXBASE%%/etc/vmware %%LINUXBASE%%/dev/vmmon %%LINUXBASE%%/dev/vmnet1 %%LINUXBASE%%/bin/df %%LINUXBASE%%/dev/tty[0-9]* %%LINUXBASE%%/dev/hd? %%LINUXBASE%%/dev/null
@dirrm etc/vmware
@dirrm lib/vmware/bin
@dirrm lib/vmware/lib/help
@dirrm lib/vmware/lib/modules
@dirrm lib/vmware/lib/xkeymap
@dirrm lib/vmware/lib/configurator
@dirrm lib/vmware/lib/isoimages
@dirrm lib/vmware/lib/licenses/site
@dirrm lib/vmware/lib/licenses/user
@dirrm lib/vmware/lib/licenses
@dirrm lib/vmware/lib/messages/chef
@dirrm lib/vmware/lib/messages/jive
@dirrm lib/vmware/lib/messages/ja
@dirrm lib/vmware/lib/messages
@dirrm lib/vmware/lib/smb/codepages
@dirrm lib/vmware/lib/smb
@dirrm lib/vmware/lib
@dirrm lib/vmware
@dirrm share/doc/vmware
