# net/3proxy/files/3proxy.sh.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: threeproxy
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

# Define these threeproxy_* variables in one of these files:
#	/etc/rc.conf
#	/etc/rc.conf.local
#	/etc/rc.conf.d/threeproxy
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE

threeproxy_enable=${threeproxy_enable-"NO"}
threeproxy_flags=${threeproxy_flags-"%%PREFIX%%/etc/3proxy.cfg"}

. /etc/rc.subr

name="threeproxy"
rcvar=threeproxy_enable
command="%%PREFIX%%/bin/3proxy"

load_rc_config $name

start_cmd="echo \"Starting ${name}.\";	${command} ${threeproxy_flags}"

run_rc_command "$1"
