# comms/conserver-com/pkg-plist
sbin/conserver
etc/rc.d/conserver.sh
bin/console
