# comms/conserver-com/pkg-plist
bin/console
etc/rc.d/conserver.sh
%%EXAMPLESDIR%%/conserver.cf
%%EXAMPLESDIR%%/conserver.passwd
%%EXAMPLESDIR%%/conserver.rc
sbin/conserver
@dirrm %%EXAMPLESDIR%%
