# net-mgmt/flowd/files/openrc-flowd.in
#!/sbin/openrc-run
#
# Add the following lines to /etc/rc.conf to run flowd:
#
# flowd_conf (file):             Local config file.
#                                Default is "%%PREFIX%%/etc/flowd.conf".
# flowd_flags (string):          Flags to pass to flowd.
#                                Default is empty.

name="flowd"

: ${flowd_conf="%%PREFIX%%/etc/flowd.conf"}

required_files="${flowd_conf}"
command="%%PREFIX%%/sbin/flowd"
command_args="-f ${flowd_conf}"
