# net/openldap24-server/pkg-plist
@comment $FreeBSD$
@unexec %%RC_DIR%%/etc/rc.d/slapd stop 2>&1 >/dev/null || true
etc/openldap/slapd.ldif.default
etc/openldap/slapd.ldif
@unexec if cmp -s %D/etc/openldap/slapd.conf %D/etc/openldap/slapd.conf.default; then rm -f %D/etc/openldap/slapd.conf; fi
etc/openldap/slapd.conf.default
@exec [ -f %B/slapd.conf ] || cp %B/%f %B/slapd.conf
etc/openldap/schema/ppolicy.ldif
etc/openldap/schema/pmi.ldif
etc/openldap/schema/openldap.ldif
etc/openldap/schema/nis.ldif
etc/openldap/schema/misc.ldif
etc/openldap/schema/java.ldif
etc/openldap/schema/inetorgperson.ldif
etc/openldap/schema/dyngroup.ldif
etc/openldap/schema/duaconf.ldif
etc/openldap/schema/cosine.ldif
etc/openldap/schema/core.ldif
etc/openldap/schema/corba.ldif
etc/openldap/schema/collective.ldif
etc/openldap/schema/README
etc/openldap/DB_CONFIG.example
%%SLAPI%%lib/libslapi.a
%%SLAPI%%lib/libslapi.la
%%SLAPI%%lib/libslapi.so
%%SLAPI%%lib/libslapi-2.4.so
%%SLAPI%%lib/libslapi-2.4.so.%%SHLIB_MAJOR%%
%%MODULES%%@exec mkdir -p %D/libexec/openldap
%%BACK_BDB%%libexec/openldap/back_bdb-2.4.so
%%BACK_BDB%%libexec/openldap/back_bdb-2.4.so.%%SHLIB_MAJOR%%
%%BACK_BDB%%libexec/openldap/back_bdb.la
%%BACK_BDB%%libexec/openldap/back_bdb.so
%%BACK_HDB%%libexec/openldap/back_hdb-2.4.so
%%BACK_HDB%%libexec/openldap/back_hdb-2.4.so.%%SHLIB_MAJOR%%
%%BACK_HDB%%libexec/openldap/back_hdb.la
%%BACK_HDB%%libexec/openldap/back_hdb.so
%%BACKEND%%libexec/openldap/back_ldap-2.4.so
%%BACKEND%%libexec/openldap/back_ldap-2.4.so.%%SHLIB_MAJOR%%
%%BACKEND%%libexec/openldap/back_ldap.la
%%BACKEND%%libexec/openldap/back_ldap.so
%%BACKEND%%libexec/openldap/back_meta-2.4.so
%%BACKEND%%libexec/openldap/back_meta-2.4.so.%%SHLIB_MAJOR%%
%%BACKEND%%libexec/openldap/back_meta.la
%%BACKEND%%libexec/openldap/back_meta.so
%%BACKEND%%libexec/openldap/back_monitor-2.4.so
%%BACKEND%%libexec/openldap/back_monitor-2.4.so.%%SHLIB_MAJOR%%
%%BACKEND%%libexec/openldap/back_monitor.la
%%BACKEND%%libexec/openldap/back_monitor.so
%%BACKEND%%libexec/openldap/back_null-2.4.so
%%BACKEND%%libexec/openldap/back_null-2.4.so.%%SHLIB_MAJOR%%
%%BACKEND%%libexec/openldap/back_null.la
%%BACKEND%%libexec/openldap/back_null.so
%%BACK_DNSSRV%%libexec/openldap/back_dnssrv-2.4.so
%%BACK_DNSSRV%%libexec/openldap/back_dnssrv-2.4.so.%%SHLIB_MAJOR%%
%%BACK_DNSSRV%%libexec/openldap/back_dnssrv.la
%%BACK_DNSSRV%%libexec/openldap/back_dnssrv.so
%%BACK_MDB%%libexec/openldap/back_mdb-2.4.so
%%BACK_MDB%%libexec/openldap/back_mdb-2.4.so.%%SHLIB_MAJOR%%
%%BACK_MDB%%libexec/openldap/back_mdb.la
%%BACK_MDB%%libexec/openldap/back_mdb.so
%%BACK_PASSWD%%libexec/openldap/back_passwd-2.4.so
%%BACK_PASSWD%%libexec/openldap/back_passwd-2.4.so.%%SHLIB_MAJOR%%
%%BACK_PASSWD%%libexec/openldap/back_passwd.la
%%BACK_PASSWD%%libexec/openldap/back_passwd.so
%%BACK_PERL%%libexec/openldap/back_perl-2.4.so
%%BACK_PERL%%libexec/openldap/back_perl-2.4.so.%%SHLIB_MAJOR%%
%%BACK_PERL%%libexec/openldap/back_perl.la
%%BACK_PERL%%libexec/openldap/back_perl.so
%%BACK_RELAY%%libexec/openldap/back_relay-2.4.so
%%BACK_RELAY%%libexec/openldap/back_relay-2.4.so.%%SHLIB_MAJOR%%
%%BACK_RELAY%%libexec/openldap/back_relay.la
%%BACK_RELAY%%libexec/openldap/back_relay.so
%%BACK_SHELL%%libexec/openldap/back_shell-2.4.so
%%BACK_SHELL%%libexec/openldap/back_shell-2.4.so.%%SHLIB_MAJOR%%
%%BACK_SHELL%%libexec/openldap/back_shell.la
%%BACK_SHELL%%libexec/openldap/back_shell.so
%%BACK_SOCK%%libexec/openldap/back_sock-2.4.so
%%BACK_SOCK%%libexec/openldap/back_sock-2.4.so.%%SHLIB_MAJOR%%
%%BACK_SOCK%%libexec/openldap/back_sock.la
%%BACK_SOCK%%libexec/openldap/back_sock.so
%%BACK_SQL%%libexec/openldap/back_sql-2.4.so
%%BACK_SQL%%libexec/openldap/back_sql-2.4.so.%%SHLIB_MAJOR%%
%%BACK_SQL%%libexec/openldap/back_sql.la
%%BACK_SQL%%libexec/openldap/back_sql.so
%%SMBPWD%%libexec/openldap/smbk5pwd.la
%%SMBPWD%%libexec/openldap/smbk5pwd.so
%%SMBPWD%%libexec/openldap/smbk5pwd.so.0
libexec/slapd
@dirrmtry libexec/openldap
sbin/slapacl
sbin/slapadd
sbin/slapauth
sbin/slapcat
sbin/slapdn
sbin/slapindex
sbin/slappasswd
sbin/slapschema
sbin/slaptest
@cwd /
@exec mkdir -p %%LDAP_RUN_DIR%%
@dirrmtry %%LDAP_RUN_DIR%%
@cwd %%PREFIX%%
