# net/openldap24-server/pkg-plist
@preunexec %%RC_DIR%%/etc/rc.d/slapd stop 2>&1 >/dev/null || true
%%ETCDIR%%/schema/README
%%ETCDIR%%/schema/collective.ldif
@sample %%ETCDIR%%/schema/collective.schema.sample
%%ETCDIR%%/schema/corba.ldif
@sample %%ETCDIR%%/schema/corba.schema.sample
%%ETCDIR%%/schema/core.ldif
@sample %%ETCDIR%%/schema/core.schema.sample
%%ETCDIR%%/schema/cosine.ldif
@sample %%ETCDIR%%/schema/cosine.schema.sample
%%ETCDIR%%/schema/duaconf.ldif
@sample %%ETCDIR%%/schema/duaconf.schema.sample
%%ETCDIR%%/schema/dyngroup.ldif
@sample %%ETCDIR%%/schema/dyngroup.schema.sample
%%ETCDIR%%/schema/inetorgperson.ldif
@sample %%ETCDIR%%/schema/inetorgperson.schema.sample
%%ETCDIR%%/schema/java.ldif
@sample %%ETCDIR%%/schema/java.schema.sample
%%ETCDIR%%/schema/misc.ldif
@sample %%ETCDIR%%/schema/misc.schema.sample
%%ETCDIR%%/schema/nis.ldif
@sample %%ETCDIR%%/schema/nis.schema.sample
%%ETCDIR%%/schema/openldap.ldif
@sample %%ETCDIR%%/schema/openldap.schema.sample
%%ETCDIR%%/schema/pmi.ldif
@sample %%ETCDIR%%/schema/pmi.schema.sample
%%ETCDIR%%/schema/ppolicy.ldif
@sample %%ETCDIR%%/schema/ppolicy.schema.sample
@sample %%ETCDIR%%/slapd.conf.sample
@sample %%ETCDIR%%/slapd.ldif.sample
%%BACK_BDB%%%%ETCDIR%%/DB_CONFIG.example
%%SLAPI%%lib/libslapi-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%SLAPI%%lib/libslapi-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%SLAPI%%lib/libslapi.a
%%SLAPI%%lib/libslapi.la
%%SLAPI%%lib/libslapi.so
%%BACK_BDB%%libexec/openldap/back_bdb-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACK_BDB%%libexec/openldap/back_bdb-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACK_BDB%%libexec/openldap/back_bdb.la
%%BACK_BDB%%libexec/openldap/back_bdb.so
%%BACK_HDB%%libexec/openldap/back_hdb-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACK_HDB%%libexec/openldap/back_hdb-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACK_HDB%%libexec/openldap/back_hdb.la
%%BACK_HDB%%libexec/openldap/back_hdb.so
%%BACKEND%%libexec/openldap/back_ldap-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACKEND%%libexec/openldap/back_ldap-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACKEND%%libexec/openldap/back_ldap.la
%%BACKEND%%libexec/openldap/back_ldap.so
%%BACKEND%%libexec/openldap/back_meta-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACKEND%%libexec/openldap/back_meta-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACKEND%%libexec/openldap/back_meta.la
%%BACKEND%%libexec/openldap/back_meta.so
%%BACKEND%%libexec/openldap/back_monitor-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACKEND%%libexec/openldap/back_monitor-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACKEND%%libexec/openldap/back_monitor.la
%%BACKEND%%libexec/openldap/back_monitor.so
%%BACKEND%%libexec/openldap/back_null-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACKEND%%libexec/openldap/back_null-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACKEND%%libexec/openldap/back_null.la
%%BACKEND%%libexec/openldap/back_null.so
%%BACK_DNSSRV%%libexec/openldap/back_dnssrv-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACK_DNSSRV%%libexec/openldap/back_dnssrv-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACK_DNSSRV%%libexec/openldap/back_dnssrv.la
%%BACK_DNSSRV%%libexec/openldap/back_dnssrv.so
%%BACK_MDB%%libexec/openldap/back_mdb-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACK_MDB%%libexec/openldap/back_mdb-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACK_MDB%%libexec/openldap/back_mdb.la
%%BACK_MDB%%libexec/openldap/back_mdb.so
%%BACK_PASSWD%%libexec/openldap/back_passwd-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACK_PASSWD%%libexec/openldap/back_passwd-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACK_PASSWD%%libexec/openldap/back_passwd.la
%%BACK_PASSWD%%libexec/openldap/back_passwd.so
%%BACK_PERL%%libexec/openldap/back_perl-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACK_PERL%%libexec/openldap/back_perl-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACK_PERL%%libexec/openldap/back_perl.la
%%BACK_PERL%%libexec/openldap/back_perl.so
%%BACK_RELAY%%libexec/openldap/back_relay-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACK_RELAY%%libexec/openldap/back_relay-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACK_RELAY%%libexec/openldap/back_relay.la
%%BACK_RELAY%%libexec/openldap/back_relay.so
%%BACK_SHELL%%libexec/openldap/back_shell-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACK_SHELL%%libexec/openldap/back_shell-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACK_SHELL%%libexec/openldap/back_shell.la
%%BACK_SHELL%%libexec/openldap/back_shell.so
%%BACK_SOCK%%libexec/openldap/back_sock-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACK_SOCK%%libexec/openldap/back_sock-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACK_SOCK%%libexec/openldap/back_sock.la
%%BACK_SOCK%%libexec/openldap/back_sock.so
%%BACK_SQL%%libexec/openldap/back_sql-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%
%%BACK_SQL%%libexec/openldap/back_sql-%%OPENLDAP_MAJOR%%.so.%%SHLIB_MAJOR%%.%%SHLIB_MINOR%%
%%BACK_SQL%%libexec/openldap/back_sql.la
%%BACK_SQL%%libexec/openldap/back_sql.so
%%SHA2%%libexec/openldap/pw-sha2.la
%%SHA2%%libexec/openldap/pw-sha2.so
%%SHA2%%libexec/openldap/pw-sha2.so.0
%%SHA2%%libexec/openldap/pw-sha2.so.0.0.0
%%SMBPWD%%libexec/openldap/smbk5pwd.la
%%SMBPWD%%libexec/openldap/smbk5pwd.so
%%SMBPWD%%libexec/openldap/smbk5pwd.so.0
%%SMBPWD%%libexec/openldap/smbk5pwd.so.0.0.0
%%LASTBIND%%libexec/openldap/lastbind.la
%%LASTBIND%%libexec/openldap/lastbind.so
%%LASTBIND%%libexec/openldap/lastbind.so.0
%%LASTBIND%%libexec/openldap/lastbind.so.0.0.0
libexec/slapd
man/man5/slapd-bdb.5.gz
man/man5/slapd-config.5.gz
man/man5/slapd-dnssrv.5.gz
man/man5/slapd-hdb.5.gz
man/man5/slapd-ldap.5.gz
man/man5/slapd-ldbm.5.gz
man/man5/slapd-ldif.5.gz
man/man5/slapd-mdb.5.gz
man/man5/slapd-meta.5.gz
man/man5/slapd-monitor.5.gz
man/man5/slapd-ndb.5.gz
man/man5/slapd-null.5.gz
man/man5/slapd-passwd.5.gz
man/man5/slapd-perl.5.gz
man/man5/slapd-relay.5.gz
man/man5/slapd-shell.5.gz
man/man5/slapd-sock.5.gz
man/man5/slapd-sql.5.gz
man/man5/slapd.access.5.gz
man/man5/slapd.backends.5.gz
man/man5/slapd.conf.5.gz
man/man5/slapd.overlays.5.gz
man/man5/slapd.plugin.5.gz
man/man5/slapo-accesslog.5.gz
man/man5/slapo-auditlog.5.gz
man/man5/slapo-chain.5.gz
man/man5/slapo-collect.5.gz
man/man5/slapo-constraint.5.gz
man/man5/slapo-dds.5.gz
man/man5/slapo-dyngroup.5.gz
man/man5/slapo-dynlist.5.gz
man/man5/slapo-memberof.5.gz
man/man5/slapo-pbind.5.gz
man/man5/slapo-pcache.5.gz
man/man5/slapo-ppolicy.5.gz
man/man5/slapo-refint.5.gz
man/man5/slapo-retcode.5.gz
man/man5/slapo-rwm.5.gz
man/man5/slapo-sock.5.gz
man/man5/slapo-sssvlv.5.gz
man/man5/slapo-syncprov.5.gz
man/man5/slapo-translucent.5.gz
man/man5/slapo-unique.5.gz
man/man5/slapo-valsort.5.gz
man/man8/slapacl.8.gz
man/man8/slapadd.8.gz
man/man8/slapauth.8.gz
man/man8/slapcat.8.gz
man/man8/slapd.8.gz
man/man8/slapdn.8.gz
man/man8/slapindex.8.gz
man/man8/slappasswd.8.gz
man/man8/slapschema.8.gz
man/man8/slaptest.8.gz
sbin/slapacl
sbin/slapadd
sbin/slapauth
sbin/slapcat
sbin/slapdn
sbin/slapindex
sbin/slappasswd
sbin/slapschema
sbin/slaptest
@dir %%ETCDIR%%/schema
%%MODULES%%@dir libexec/openldap
@dir(%%LDAP_USER%%,%%LDAP_GROUP%%,775) %%LDAP_RUN_DIR%%
@dir(%%LDAP_USER%%,%%LDAP_GROUP%%,775) %%DATABASEDIR%%
