# net/syncthing-discosrv/files/syncthing-discosrv.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: syncthingdiscosrv
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# syncthingdiscosrv_enable (bool):	Set to NO by default.
#				Set it to YES to enable syncthing-discosrv.
# syncthingdiscosrv_user (user):        Set user to run syncthing-discosrv.
#                               Default is "syncthing".
# syncthingdiscosrv_group (group):      Set group to run syncthing-discosrv.
#                               Default is "syncthing".
# syncthingdiscosrv_dir (dir):          Set dir to run syncthing-discosrv in.
#                               Default is "/var/db/syncthing-discosrv".

. /etc/rc.subr

name=syncthingdiscosrv
rcvar=syncthingdiscosrv_enable

load_rc_config $name

: ${syncthingdiscosrv_enable:="NO"}
: ${syncthingdiscosrv_user:="syncthing"}
: ${syncthingdiscosrv_group:="syncthing"}
: ${syncthingdiscosrv_dir:="/var/db/syncthing-discosrv"}

export STNORESTART=true

pidfile=/var/run/syncthingdiscosrv.pid
procname="%%PREFIX%%/bin/syncthing-discosrv"
command="/usr/sbin/daemon"
command_args="-f -p ${pidfile} ${procname} -db-dir=${syncthingdiscosrv_dir}/db -stats-file=${syncthingdiscosrv_dir}/stats ${syncthingdiscosrv_flags}"

start_precmd=syncthingdiscosrv_startprecmd

syncthingdiscosrv_startprecmd()
{
        if [ ! -e ${pidfile} ]; then
                install -o ${syncthingdiscosrv_user} -g ${syncthingdiscosrv_group} /dev/null ${pidfile};
        fi

        if [ ! -d ${syncthingdiscosrv_dir} ]; then
                install -d -o ${syncthingdiscosrv_user} -g ${syncthingdiscosrv_group} ${syncthingdiscosrv_dir}
        fi

}

run_rc_command "$1"
