# security/razorback-dispatcher/files/openrc-dispatcher_safed.in
#!/sbin/openrc-run

name=dispatcher_safed

pidfile="/var/run/dispatcher_safed.pid"
command="%%PREFIX%%/bin/dispatcher_safed"
supervisor=supervise-daemon

depend() {
	need dispatcherd
	keyword -shutdown
}
