# security/silktools/pkg-plist
bin/mapsid
bin/num2dot
bin/rwaddrcount
bin/rwaggbag
bin/rwaggbagbuild
bin/rwaggbagcat
bin/rwaggbagtool
bin/rwallformats
bin/rwappend
bin/rwbag
bin/rwbagbuild
bin/rwbagcat
bin/rwbagtool
bin/rwcat
bin/rwcombine
bin/rwcompare
bin/rwcount
bin/rwcut
bin/rwdedupe
bin/rwfglob
bin/rwfileinfo
bin/rwfilter
bin/rwgeoip2ccmap
bin/rwgroup
bin/rwidsquery
bin/rwip2cc
bin/rwipfix2silk
bin/rwmatch
bin/rwnetmask
bin/rwp2yaf2silk
bin/rwpcut
bin/rwpdedupe
bin/rwpdu2silk
bin/rwpmapbuild
bin/rwpmapcat
bin/rwpmaplookup
bin/rwpmatch
bin/rwptoflow
bin/rwrandomizeip
bin/rwrecgenerator
bin/rwresolve
bin/rwscan
bin/rwscanquery
bin/rwset
bin/rwsetbuild
bin/rwsetcat
bin/rwsetmember
bin/rwsettool
bin/rwsilk2ipfix
bin/rwsiteinfo
bin/rwsort
bin/rwsplit
bin/rwstats
bin/rwswapbytes
bin/rwtotal
bin/rwtuc
bin/rwuniq
bin/silk_config
include/silk/bagtree.h
include/silk/hashlib.h
include/silk/iptree.h
include/silk/libflowsource.h
include/silk/probeconf.h
include/silk/redblack.h
include/silk/rwascii.h
include/silk/rwflowpack.h
include/silk/rwpack.h
include/silk/rwrec.h
include/silk/silk.h
include/silk/silk_config.h
include/silk/silk_config2.h
include/silk/silk_files.h
include/silk/silk_types.h
include/silk/skaggbag.h
include/silk/skbag.h
include/silk/skcountry.h
include/silk/skdaemon.h
include/silk/skdeque.h
include/silk/skdllist.h
include/silk/skheader.h
include/silk/skheap.h
include/silk/skipaddr.h
include/silk/skipset.h
include/silk/sklog.h
include/silk/skmempool.h
include/silk/sknetstruct.h
include/silk/skplugin.h
include/silk/skpolldir.h
include/silk/skprefixmap.h
include/silk/sksite.h
include/silk/skstream.h
include/silk/skstringmap.h
include/silk/sktempfile.h
include/silk/skthread.h
include/silk/sktimer.h
include/silk/sktracemsg.h
include/silk/skvector.h
include/silk/utils.h
%%STATIC%%lib/libflowsource.a
lib/libflowsource.so
lib/libflowsource.so.17
lib/libflowsource.so.17.0.0
%%STATIC%%lib/libsilk-thrd.a
lib/libsilk-thrd.so
lib/libsilk-thrd.so.5
lib/libsilk-thrd.so.5.1.4
%%STATIC%%lib/libsilk.a
lib/libsilk.so
lib/libsilk.so.25
lib/libsilk.so.25.0.0
%%STATIC%%lib/silk/app-mismatch.a
lib/silk/app-mismatch.so
%%STATIC%%lib/silk/conficker-c.a
lib/silk/conficker-c.so
%%STATIC%%lib/silk/cutmatch.a
lib/silk/cutmatch.so
%%STATIC%%lib/silk/flowkey.a
lib/silk/flowkey.so
%%STATIC%%lib/silk/flowrate.a
lib/silk/flowrate.so
%%STATIC%%lib/silk/int-ext-fields.a
lib/silk/int-ext-fields.so
%%STATIC%%lib/silk/packlogic-generic.a
lib/silk/packlogic-generic.so
%%STATIC%%lib/silk/packlogic-twoway.a
lib/silk/packlogic-twoway.so
man/man1/mapsid.1.gz
man/man1/num2dot.1.gz
man/man1/rwaddrcount.1.gz
man/man1/rwaggbag.1.gz
man/man1/rwaggbagbuild.1.gz
man/man1/rwaggbagcat.1.gz
man/man1/rwaggbagtool.1.gz
man/man1/rwappend.1.gz
man/man1/rwbag.1.gz
man/man1/rwbagbuild.1.gz
man/man1/rwbagcat.1.gz
man/man1/rwbagtool.1.gz
man/man1/rwcat.1.gz
man/man1/rwcombine.1.gz
man/man1/rwcompare.1.gz
man/man1/rwcount.1.gz
man/man1/rwcut.1.gz
man/man1/rwdedupe.1.gz
man/man1/rwfglob.1.gz
man/man1/rwfileinfo.1.gz
man/man1/rwfilter.1.gz
man/man1/rwgeoip2ccmap.1.gz
man/man1/rwgroup.1.gz
man/man1/rwidsquery.1.gz
man/man1/rwip2cc.1.gz
man/man1/rwipfix2silk.1.gz
man/man1/rwmatch.1.gz
man/man1/rwnetmask.1.gz
man/man1/rwp2yaf2silk.1.gz
man/man1/rwpcut.1.gz
man/man1/rwpdedupe.1.gz
man/man1/rwpdu2silk.1.gz
man/man1/rwpmapbuild.1.gz
man/man1/rwpmapcat.1.gz
man/man1/rwpmaplookup.1.gz
man/man1/rwpmatch.1.gz
man/man1/rwptoflow.1.gz
man/man1/rwrandomizeip.1.gz
man/man1/rwrecgenerator.1.gz
man/man1/rwresolve.1.gz
man/man1/rwscan.1.gz
man/man1/rwscanquery.1.gz
man/man1/rwset.1.gz
man/man1/rwsetbuild.1.gz
man/man1/rwsetcat.1.gz
man/man1/rwsetmember.1.gz
man/man1/rwsettool.1.gz
man/man1/rwsilk2ipfix.1.gz
man/man1/rwsiteinfo.1.gz
man/man1/rwsort.1.gz
man/man1/rwsplit.1.gz
man/man1/rwstats.1.gz
man/man1/rwswapbytes.1.gz
man/man1/rwtotal.1.gz
man/man1/rwtuc.1.gz
man/man1/rwuniq.1.gz
man/man1/silk_config.1.gz
man/man3/addrtype.3.gz
man/man3/app-mismatch.3.gz
man/man3/ccfilter.3.gz
man/man3/conficker-c.3.gz
man/man3/cutmatch.3.gz
man/man3/flowkey.3.gz
man/man3/flowrate.3.gz
man/man3/int-ext-fields.3.gz
man/man3/packlogic-generic.3.gz
man/man3/packlogic-twoway.3.gz
man/man3/pmapfilter.3.gz
%%PYTHON%%man/man3/pysilk.3.gz
man/man3/silk-plugin.3.gz
%%PYTHON%%man/man3/silkpython.3.gz
man/man5/sensor.conf.5.gz
man/man5/silk.conf.5.gz
man/man7/silk.7.gz
man/man8/flowcap.8.gz
man/man8/rwflowappend.8.gz
man/man8/rwflowpack.8.gz
man/man8/rwguess.8.gz
man/man8/rwpackchecker.8.gz
man/man8/rwpollexec.8.gz
man/man8/rwreceiver.8.gz
man/man8/rwsender.8.gz
sbin/flowcap
sbin/rwflowappend
sbin/rwflowpack
sbin/rwguess
sbin/rwpackchecker
sbin/rwpollexec
sbin/rwreceiver
sbin/rwsender
share/silk/addrtype-templ.txt
share/silk/etc/flowcap.conf
share/silk/etc/init.d/flowcap
share/silk/etc/init.d/rwflowappend
share/silk/etc/init.d/rwflowpack
share/silk/etc/init.d/rwpollexec
share/silk/etc/init.d/rwreceiver
share/silk/etc/init.d/rwsender
share/silk/etc/rwflowappend.conf
share/silk/etc/rwflowpack.conf
share/silk/etc/rwpollexec.conf
share/silk/etc/rwreceiver.conf
share/silk/etc/rwsender.conf
share/silk/generic-silk.conf
share/silk/silk.magic
share/silk/twoway-silk.conf
%%PYTHON%%%%PYTHON_SITELIBDIR%%/netsa_silk.py
%%PYTHON%%%%PYTHON_SITELIBDIR%%/netsa_silk.pyc
%%PYTHON%%%%PYTHON_SITELIBDIR%%/netsa_silk.pyo
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/__init__.py
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/__init__.pyc
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/__init__.pyo
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/_netsa_silk.py
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/_netsa_silk.pyc
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/_netsa_silk.pyo
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/fglob.py
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/fglob.pyc
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/fglob.pyo
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/plugin.py
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/plugin.pyc
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/plugin.pyo
%%STATIC%%%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/pysilk.a
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/pysilk.so
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/site.py
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/site.pyc
%%PYTHON%%%%PYTHON_SITELIBDIR%%/silk/site.pyo
