# security/silktools/pkg-plist
@comment $FreeBSD$
bin/rwdedupe
bin/buildset
bin/mapsid
bin/num2dot
bin/readset
bin/rwaddrcount
bin/rwappend
bin/rwbag
bin/rwbagbuild
bin/rwbagcat
bin/rwbagtool
bin/rwcat
bin/rwcount
bin/rwcut
bin/rwip2cc
bin/rwpmapbuild
bin/rwpmapcat
bin/rwallformats
bin/rwresolve
bin/rwsetbuild
bin/rwsetcat
bin/rwsetintersect
bin/rwsetunion
bin/rwtuc
bin/rwp2yaf2silk
bin/rwfglob
bin/rwfileinfo
bin/rwfilter
bin/rwgroup
bin/rwmatch
bin/rwnetmask
bin/rwpcut
bin/rwpdedupe
bin/rwgeoip2ccmap
bin/rwpmatch
bin/rwptoflow
bin/rwrandomizeip
bin/rwrtd2split
bin/rwscan
bin/rwscanquery
bin/rwset
bin/rwsetmember
bin/rwsettool
bin/rwset-union
bin/rwsort
bin/rwstats
bin/rwsplit
bin/rwswapbytes
bin/rwtotal
bin/rwuniq
bin/setintersect
lib/libsilk.so.0
lib/libsilk.so
lib/libsilk.la
lib/libsilk.a
lib/libthreadutils.so.0
lib/libthreadutils.so
lib/libthreadutils.la
lib/libthreadutils.a
lib/libredblack.so.0
lib/libredblack.so
lib/libredblack.la
lib/libredblack.a
lib/libsksetbag.so.0
lib/libsksetbag.so
lib/libsksetbag.la
lib/libsksetbag.a
lib/silk/libaddrtype.so.0
lib/silk/libaddrtype.so
lib/silk/libaddrtype.la
lib/silk/libaddrtype.a
lib/silk/libccfilter.so.0
lib/silk/libccfilter.so
lib/silk/libccfilter.la
lib/silk/libccfilter.a
lib/silk/libpmapfilter.so.0
lib/silk/libpmapfilter.so
lib/silk/libpmapfilter.la
lib/silk/libpmapfilter.a
lib/silk/libipport.so.0
lib/silk/libipport.so
lib/silk/libipport.la
lib/silk/libipport.a
lib/silk/libippair.so.0
lib/silk/libippair.so
lib/silk/libippair.la
lib/silk/libippair.a
lib/silk/librwslammer.so.0
lib/silk/librwslammer.so
lib/silk/librwslammer.la
lib/silk/librwslammer.a
lib/silk/librwsynackfin.so.0
lib/silk/librwsynackfin.so
lib/silk/librwsynackfin.la
lib/silk/librwsynackfin.a
lib/silk/librwp2f_minbytes.so.0
lib/silk/librwp2f_minbytes.so
lib/silk/librwp2f_minbytes.la
lib/silk/librwp2f_minbytes.a
lib/libhash.so.0
lib/libhash.so
lib/libhash.la
lib/libhash.a
lib/libheap.so.0
lib/libheap.so
lib/libheap.la
lib/libheap.a
lib/libflowsource.so.0
lib/libflowsource.so
lib/libflowsource.la
lib/libflowsource.a
lib/libskdeque.so.0
lib/libskdeque.so
lib/libskdeque.la
lib/libskdeque.a
lib/libskpolldir.so.0
lib/libskpolldir.so
lib/libskpolldir.la
lib/libskpolldir.a
lib/libskmsg.so.0
lib/libskmsg.so
lib/libskmsg.la
lib/libskmsg.a
lib/libsktransfer.so.0
lib/libsktransfer.so
lib/libsktransfer.la
lib/libsktransfer.a
sbin/flowcap
sbin/rwflowpack
sbin/rwpackchecker
sbin/rwflowappend
sbin/rwreceiver
sbin/rwsender
share/silk/etc/flowcap.conf
share/silk/etc/rwflowpack.conf
share/silk/etc/init.d/flowcap
share/silk/etc/init.d/rwflowpack
share/silk/etc/init.d/rwflowappend
share/silk/etc/init.d/rwreceiver
share/silk/etc/init.d/rwsender
share/silk/etc/rwflowappend.conf
share/silk/etc/rwreceiver.conf
share/silk/etc/rwsender.conf
share/silk/silk.conf.template
include/silk/bagtree.h
include/silk/convert-cflowd.h
include/silk/convert-flowcap.h
include/silk/dynlib.h
include/silk/flowcap-source.h
include/silk/flowcap-stats.h
include/silk/flowcap.h
include/silk/gnu_getopt.h
include/silk/hashlib.h
include/silk/hashlib_conf.h
include/silk/hashlib_utils.h
include/silk/hashwrap.h
include/silk/heaplib.h
include/silk/iochecks.h
include/silk/ipfixsource.h
include/silk/iptree.h
include/silk/libflowsource.h
include/silk/libsksetbag.h
include/silk/lzo-file.h
include/silk/multiqueue.h
include/silk/pdusource.h
include/silk/probeconf.h
include/silk/redblack.h
include/silk/rwpack.h
include/silk/silk.h
include/silk/silk_files.h
include/silk/silk_site.h
include/silk/skdaemon.h
include/silk/skdeque.h
include/silk/skdllist.h
include/silk/skiobuf.h
include/silk/skipfix.h
include/silk/sklinkedlist.h
include/silk/sklog.h
include/silk/skmsg.h
include/silk/skpolldir.h
include/silk/skprefixmap.h
include/silk/skprintnets.h
include/silk/sksite.h
include/silk/skstream.h
include/silk/skstringmap.h
include/silk/skthread.h
include/silk/sktimer.h
include/silk/sktransfer.h
include/silk/skvector.h
include/silk/utils.h
include/silk/v5pdu.h
include/silk/silk_site_twoway.h
@dirrmtry share/silk/etc/init.d
@dirrmtry share/silk/etc
@dirrmtry share/silk
@dirrmtry lib/silk
@dirrmtry include/silk
