# security/silktools/pkg-plist
bin/mapsid
bin/num2dot
bin/rwaddrcount
bin/rwallformats
bin/rwappend
bin/rwbag
bin/rwbagbuild
bin/rwbagcat
bin/rwbagtool
bin/rwcat
bin/rwcount
bin/rwcut
bin/rwdedupe
bin/rwfglob
bin/rwfileinfo
bin/rwfilter
bin/rwgeoip2ccmap
bin/rwgroup
bin/rwidsquery
bin/rwip2cc
bin/rwipfix2silk
bin/rwmatch
bin/rwnetmask
bin/rwp2yaf2silk
bin/rwpcut
bin/rwpdedupe
bin/rwpmapbuild
bin/rwpmapcat
bin/rwpmatch
bin/rwptoflow
bin/rwrandomizeip
bin/rwresolve
bin/rwrtd2split
bin/rwscan
bin/rwscanquery
bin/rwset
bin/rwsetbuild
bin/rwsetcat
bin/rwsetintersect
bin/rwsetmember
bin/rwsettool
bin/rwsetunion
bin/rwsilk2ipfix
bin/rwsort
bin/rwsplit
bin/rwstats
bin/rwswapbytes
bin/rwtotal
bin/rwtuc
bin/rwuniq
include/silk/bagtree.h
include/silk/convert-cflowd.h
include/silk/dynlib.h
include/silk/flowcap-source.h
include/silk/gnu_getopt.h
include/silk/hashlib.h
include/silk/heaplib.h
include/silk/iochecks.h
include/silk/ipfixsource.h
include/silk/iptree.h
include/silk/libflowsource.h
include/silk/libsksetbag.h
include/silk/multiqueue.h
include/silk/pdusource.h
include/silk/probeconf.h
include/silk/redblack.h
include/silk/rwflowpack.h
include/silk/rwpack.h
include/silk/rwrec.h
include/silk/silk.h
include/silk/silk_config.h
include/silk/silk_files.h
include/silk/skdaemon.h
include/silk/skdeque.h
include/silk/skdllist.h
include/silk/skheader.h
include/silk/skiobuf.h
include/silk/skipfix.h
include/silk/sklinkedlist.h
include/silk/sklog.h
include/silk/skmsg.h
include/silk/skpolldir.h
include/silk/skprefixmap.h
include/silk/skprintnets.h
include/silk/sksite.h
include/silk/skstream.h
include/silk/skstringmap.h
include/silk/skthread.h
include/silk/sktimer.h
include/silk/sktransfer.h
include/silk/skvector.h
include/silk/utils.h
include/silk/v5pdu.h
lib/libflowsource.a
lib/libflowsource.la
lib/libflowsource.so
lib/libflowsource.so.1
lib/libhash.a
lib/libhash.la
lib/libhash.so
lib/libhash.so.1
lib/libheap.a
lib/libheap.la
lib/libheap.so
lib/libheap.so.1
lib/libsilk-thrd.a
lib/libsilk-thrd.la
lib/libsilk-thrd.so
lib/libsilk-thrd.so.1
lib/libsilk.a
lib/libsilk.la
lib/libsilk.so
lib/libsilk.so.1
lib/libskipfix.a
lib/libskipfix.la
lib/libskipfix.so
lib/libskipfix.so.1
lib/libskmsg.a
lib/libskmsg.la
lib/libskmsg.so
lib/libskmsg.so.1
lib/silk/addrtype.a
lib/silk/addrtype.la
lib/silk/addrtype.so
lib/silk/ccfilter.a
lib/silk/ccfilter.la
lib/silk/ccfilter.so
lib/silk/cutmatch.a
lib/silk/cutmatch.la
lib/silk/cutmatch.so
lib/silk/flowrate.a
lib/silk/flowrate.la
lib/silk/flowrate.so
lib/silk/packlogic-generic.a
lib/silk/packlogic-generic.la
lib/silk/packlogic-generic.so
lib/silk/packlogic-twoway.a
lib/silk/packlogic-twoway.la
lib/silk/packlogic-twoway.so
lib/silk/pmapfilter.a
lib/silk/pmapfilter.la
lib/silk/pmapfilter.so
sbin/flowcap
sbin/rwflowappend
sbin/rwflowpack
sbin/rwguess
sbin/rwpackchecker
sbin/rwreceiver
sbin/rwsender
share/silk/addrtype-templ.txt
share/silk/etc/flowcap.conf
share/silk/etc/init.d/flowcap
share/silk/etc/init.d/rwflowappend
share/silk/etc/init.d/rwflowpack
share/silk/etc/init.d/rwreceiver
share/silk/etc/init.d/rwsender
share/silk/etc/rwflowappend.conf
share/silk/etc/rwflowpack.conf
share/silk/etc/rwreceiver.conf
share/silk/etc/rwsender.conf
share/silk/generic-silk.conf
share/silk/twoway-silk.conf
share/silk/update-sensor-conf
@dirrmtry share/silk/etc/init.d
@dirrmtry share/silk/etc
@dirrmtry share/silk
@dirrmtry lib/silk
@dirrmtry include/silk
