# net/minisapserver/files/openrc-sapserver.in
#!/sbin/openrc-run

name=sapserver
command=%%PREFIX%%/bin/${name}
required_files=%%PREFIX%%/etc/sap.cfg

depend() {
	use net
	keyword -shutdown
}
