# net-mgmt/nrpe/pkg-plist
sbin/nrpe
etc/nrpe.cfg-sample
libexec/nagios/check_nrpe
etc/rc.d/nrpe.sh
