# multimedia/ptx-kmod/files/openrc-ptx.in
#!/sbin/openrc-run

name="ptx"

start()
{
	load_kld %%KMODDIR%%/ptx.ko
}
stop()
{
	/sbin/kldunload ptx
}

depend()
{
	keyword -jail
}
