# multimedia/ptx-kmod/files/ptx.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: ptx
# REQUIRE: LOGIN abi
# KEYWORD: nojail

. /etc/rc.subr

name="ptx"
rcvar=ptx_enable
start_cmd="ptx_start"
stop_cmd="ptx_stop"

ptx_start()
{
	load_kld %%KMODDIR%%/ptx.ko
}
ptx_stop()
{
	/sbin/kldunload ptx
}

: ${ptx_enable="YES"}
load_rc_config $name
run_rc_command "$1"
