# x11/login.app/pkg-plist
etc/rc.d/loginapp.sh
GNUstep/Apps/Login.app/Login
GNUstep/Apps/Login.app/Preferences
GNUstep/Apps/Login.app/Scripts/Login.root
GNUstep/Apps/Login.app/Scripts/Login.user
GNUstep/Apps/Login.app/Scripts/Logout.root
GNUstep/Apps/Login.app/Scripts/PowerOff
GNUstep/Apps/Login.app/Scripts/Restart
GNUstep/Apps/Login.app/Scripts/Startup
GNUstep/Apps/Login.app/Scripts/Suspend
GNUstep/Apps/Login.app/Themes/Drone/Settings
GNUstep/Apps/Login.app/Themes/Drone/fieldlabels.png
GNUstep/Apps/Login.app/Themes/Drone/inputfield1.png
GNUstep/Apps/Login.app/Themes/Drone/inputfield2.png
GNUstep/Apps/Login.app/Themes/Drone/panel.png
GNUstep/Apps/Login.app/Themes/Drone/powerbuttondown.png
GNUstep/Apps/Login.app/Themes/Drone/powerbuttonup.png
GNUstep/Apps/Login.app/Themes/Drone/restartbuttondown.png
GNUstep/Apps/Login.app/Themes/Drone/restartbuttonup.png
GNUstep/Apps/Login.app/Themes/Drone/stop.png
GNUstep/Apps/Login.app/Themes/GNUstep/Settings
GNUstep/Apps/Login.app/Themes/GNUstep/fieldlabels.tiff
GNUstep/Apps/Login.app/Themes/GNUstep/inputfield.tiff
GNUstep/Apps/Login.app/Themes/GNUstep/panel.tiff
GNUstep/Apps/Login.app/Themes/GNUstep/powerbuttondown.tiff
GNUstep/Apps/Login.app/Themes/GNUstep/powerbuttonup.tiff
GNUstep/Apps/Login.app/Themes/GNUstep/restartbuttondown.tiff
GNUstep/Apps/Login.app/Themes/GNUstep/restartbuttonup.tiff
GNUstep/Apps/Login.app/Themes/GNUstep/stop.tiff
GNUstep/Apps/Login.app/Themes/Linux/Settings
GNUstep/Apps/Login.app/Themes/Linux/fieldlabels.tiff
GNUstep/Apps/Login.app/Themes/Linux/inputfield.tiff
GNUstep/Apps/Login.app/Themes/Linux/panel.tiff
GNUstep/Apps/Login.app/Themes/Linux/powerbuttondown.tiff
GNUstep/Apps/Login.app/Themes/Linux/powerbuttonup.tiff
GNUstep/Apps/Login.app/Themes/Linux/restartbuttondown.tiff
GNUstep/Apps/Login.app/Themes/Linux/restartbuttonup.tiff
GNUstep/Apps/Login.app/Themes/Linux/stop.tiff
@dirrm GNUstep/Apps/Login.app/Themes/Linux
@dirrm GNUstep/Apps/Login.app/Themes/GNUstep
@dirrm GNUstep/Apps/Login.app/Themes/Drone
@dirrm GNUstep/Apps/Login.app/Themes
@dirrm GNUstep/Apps/Login.app/Scripts
@dirrm GNUstep/Apps/Login.app
@dirrm GNUstep/Apps
