# hardenedbsd/secadm-devel/pkg-plist
/%%KMODDIR%%/secadm.ko
etc/rc.d/secadm
include/libsecadm.h
include/secadm.h
lib/libsecadm.so
lib/libsecadm.so.0
sbin/secadm
