# news/inn-current/pkg-plist
etc/rc.d/innd.sh
news/bin/actmerge
news/bin/actsync
news/bin/actsyncd
news/bin/archive
news/bin/auth/passwd/auth_smb
news/bin/auth/passwd/ckpasswd
news/bin/auth/passwd/radius
news/bin/auth/resolv/domain
news/bin/auth/resolv/ident
news/bin/batcher
news/bin/buffchan
news/bin/buffindexed_d
news/bin/cnfsheadconf
news/bin/cnfsstat
news/bin/control/checkgroups.pl
news/bin/control/ihave.pl
news/bin/control/newgroup.pl
news/bin/control/rmgroup.pl
news/bin/control/sendme.pl
news/bin/control/sendsys.pl
news/bin/control/senduuname.pl
news/bin/control/version.pl
news/bin/controlbatch
news/bin/controlchan
news/bin/convdate
news/bin/ctlinnd
news/bin/cvtbatch
news/bin/docheckgroups
news/bin/expire
news/bin/expireover
news/bin/expirerm
news/bin/fastrm
news/bin/filechan
news/bin/filter/INN.py
news/bin/filter/filter.tcl
news/bin/filter/filter_innd.pl
news/bin/filter/filter_innd.py
news/bin/filter/filter_nnrpd.pl
news/bin/filter/nnrpd_auth.pl
news/bin/filter/nnrpd_auth.py
news/bin/filter/startup.tcl
news/bin/filter/startup_innd.pl
news/bin/getlist
news/bin/gpgverify
news/bin/grephistory
news/bin/imapfeed
news/bin/inews
news/bin/innbind
news/bin/inncheck
news/bin/innconfval
news/bin/innd
news/bin/inndf
news/bin/innfeed
news/bin/innmail
news/bin/innreport
news/bin/innstat
news/bin/innupgrade
news/bin/innwatch
news/bin/innxbatch
news/bin/innxmit
news/bin/mailpost
news/bin/makedbz
news/bin/makehistory
news/bin/mod-active
news/bin/news.daily
news/bin/news2mail
news/bin/ninpaths
news/bin/nnrpd
news/bin/nntpget
news/bin/nntpsend
news/bin/ovdb_init
news/bin/ovdb_monitor
news/bin/ovdb_server
news/bin/ovdb_stat
news/bin/overchan
news/bin/perl-nocem
news/bin/pgpverify
news/bin/procbatch
news/bin/prunehistory
news/bin/pullnews
news/bin/rc.news
news/bin/rnews
news/bin/rnews.libexec/c7unbatch
news/bin/rnews.libexec/decode
news/bin/rnews.libexec/encode
news/bin/rnews.libexec/gunbatch
news/bin/scanlogs
news/bin/scanspool
news/bin/send-ihave
news/bin/send-nntp
news/bin/send-uucp
news/bin/sendinpaths
news/bin/sendxbatches
news/bin/shlock
news/bin/shrinkfile
news/bin/signcontrol
news/bin/simpleftp
news/bin/sm
news/bin/tally.control
news/bin/tdx-util
news/bin/writelog
news/db/active
news/db/active.times
news/db/history
news/db/history.dir
%%SUB_WITHOUT_TAGGED_HASH%%news/db/history.hash
%%SUB_WITHOUT_TAGGED_HASH%%news/db/history.index
%%SUB_WITH_TAGGED_HASH%%news/db/history.pag
news/db/newsgroups
news/doc/CONTRIBUTORS
news/doc/GPL
news/doc/HACKING
news/doc/INSTALL
news/doc/LICENSE
news/doc/NEWS
news/doc/README
news/doc/README.snapshot
news/doc/TODO
news/doc/compliance-nntp
news/doc/config-design
news/doc/config-semantics
news/doc/config-syntax
news/doc/external-auth
news/doc/history
news/doc/hook-perl
news/doc/hook-python
news/doc/hook-tcl
news/doc/sample-control
@unexec if cmp -s %D/%%ETCDIR%%/actsync.cfg %D/%%ETCDIR%%/actsync.cfg.dist; then rm -f %D/%%ETCDIR%%/actsync.cfg ; fi
%%ETCDIR%%/actsync.cfg.dist
@unexec if cmp -s %D/%%ETCDIR%%/actsync.ign %D/%%ETCDIR%%/actsync.ign.dist; then rm -f %D/%%ETCDIR%%/actsync.ign ; fi
%%ETCDIR%%/actsync.ign.dist
@unexec if cmp -s %D/%%ETCDIR%%/buffindexed.conf %D/%%ETCDIR%%/buffindexed.conf.dist; then rm -f %D/%%ETCDIR%%/buffindexed.conf ; fi
%%ETCDIR%%/buffindexed.conf.dist
@unexec if cmp -s %D/%%ETCDIR%%/control.ctl %D/%%ETCDIR%%/control.ctl.dist; then rm -f %D/%%ETCDIR%%/control.ctl ; fi
%%ETCDIR%%/control.ctl.dist
@unexec if cmp -s %D/%%ETCDIR%%/cycbuff.conf %D/%%ETCDIR%%/cycbuff.conf.dist; then rm -f %D/%%ETCDIR%%/cycbuff.conf ; fi
%%ETCDIR%%/cycbuff.conf.dist
@unexec if cmp -s %D/%%ETCDIR%%/distrib.pats %D/%%ETCDIR%%/distrib.pats.dist; then rm -f %D/%%ETCDIR%%/distrib.pats ; fi
%%ETCDIR%%/distrib.pats.dist
@unexec if cmp -s %D/%%ETCDIR%%/expire.ctl %D/%%ETCDIR%%/expire.ctl.dist; then rm -f %D/%%ETCDIR%%/expire.ctl ; fi
%%ETCDIR%%/expire.ctl.dist
@unexec if cmp -s %D/%%ETCDIR%%/incoming.conf %D/%%ETCDIR%%/incoming.conf.dist; then rm -f %D/%%ETCDIR%%/incoming.conf ; fi
%%ETCDIR%%/incoming.conf.dist
@unexec if cmp -s %D/%%ETCDIR%%/inn.conf %D/%%ETCDIR%%/inn.conf.dist; then rm -f %D/%%ETCDIR%%/inn.conf ; fi
%%ETCDIR%%/inn.conf.dist
@unexec if cmp -s %D/%%ETCDIR%%/innfeed.conf %D/%%ETCDIR%%/innfeed.conf.dist; then rm -f %D/%%ETCDIR%%/innfeed.conf ; fi
%%ETCDIR%%/innfeed.conf.dist
@unexec if cmp -s %D/%%ETCDIR%%/innreport.conf %D/%%ETCDIR%%/innreport.conf.dist; then rm -f %D/%%ETCDIR%%/innreport.conf ; fi
%%ETCDIR%%/innreport.conf.dist
@unexec if cmp -s %D/%%ETCDIR%%/innwatch.ctl %D/%%ETCDIR%%/innwatch.ctl.dist; then rm -f %D/%%ETCDIR%%/innwatch.ctl ; fi
%%ETCDIR%%/innwatch.ctl.dist
@unexec if cmp -s %D/%%ETCDIR%%/moderators %D/%%ETCDIR%%/moderators.dist; then rm -f %D/%%ETCDIR%%/moderators ; fi
%%ETCDIR%%/moderators.dist
@unexec if cmp -s %D/%%ETCDIR%%/motd.news %D/%%ETCDIR%%/motd.news.dist; then rm -f %D/%%ETCDIR%%/motd.news ; fi
%%ETCDIR%%/motd.news.dist
@unexec if cmp -s %D/%%ETCDIR%%/news2mail.cf %D/%%ETCDIR%%/news2mail.cf.dist; then rm -f %D/%%ETCDIR%%/news2mail.cf ; fi
%%ETCDIR%%/news2mail.cf.dist
@unexec if cmp -s %D/%%ETCDIR%%/newsfeeds %D/%%ETCDIR%%/newsfeeds.dist; then rm -f %D/%%ETCDIR%%/newsfeeds ; fi
%%ETCDIR%%/newsfeeds.dist
@unexec if cmp -s %D/%%ETCDIR%%/nnrpd.track %D/%%ETCDIR%%/nnrpd.track.dist; then rm -f %D/%%ETCDIR%%/nnrpd.track ; fi
%%ETCDIR%%/nnrpd.track.dist
@unexec if cmp -s %D/%%ETCDIR%%/nntpsend.ctl %D/%%ETCDIR%%/nntpsend.ctl.dist; then rm -f %D/%%ETCDIR%%/nntpsend.ctl ; fi
%%ETCDIR%%/nntpsend.ctl.dist
@unexec if cmp -s %D/%%ETCDIR%%/ovdb.conf %D/%%ETCDIR%%/ovdb.conf.dist; then rm -f %D/%%ETCDIR%%/ovdb.conf ; fi
%%ETCDIR%%/ovdb.conf.dist
@unexec if cmp -s %D/%%ETCDIR%%/overview.fmt %D/%%ETCDIR%%/overview.fmt.dist; then rm -f %D/%%ETCDIR%%/overview.fmt ; fi
%%ETCDIR%%/overview.fmt.dist
@unexec if cmp -s %D/%%ETCDIR%%/passwd.nntp %D/%%ETCDIR%%/passwd.nntp.dist; then rm -f %D/%%ETCDIR%%/passwd.nntp ; fi
%%ETCDIR%%/passwd.nntp.dist
@unexec if cmp -s %D/%%ETCDIR%%/radius.conf %D/%%ETCDIR%%/radius.conf.dist; then rm -f %D/%%ETCDIR%%/radius.conf ; fi
%%ETCDIR%%/radius.conf.dist
@unexec if cmp -s %D/%%ETCDIR%%/readers.conf %D/%%ETCDIR%%/readers.conf.dist; then rm -f %D/%%ETCDIR%%/readers.conf; fi
%%ETCDIR%%/readers.conf.dist
@unexec if cmp -s %D/%%ETCDIR%%/sasl.conf %D/%%ETCDIR%%/sasl.conf.dist; then rm -f %D/%%ETCDIR%%/sasl.conf ; fi
%%ETCDIR%%/sasl.conf.dist
@unexec if cmp -s %D/%%ETCDIR%%/storage.conf %D/%%ETCDIR%%/storage.conf.dist; then rm -f %D/%%ETCDIR%%/storage.conf ; fi
%%ETCDIR%%/storage.conf.dist
@unexec if cmp -s %D/%%ETCDIR%%/subscriptions %D/%%ETCDIR%%/subscriptions.dist; then rm -f %D/%%ETCDIR%%/subscriptions ; fi
%%ETCDIR%%/subscriptions.dist
news/include/clibrary.h
news/include/conffile.h
news/include/config.h
news/include/dbz.h
news/include/inn/buffer.h
news/include/inn/confparse.h
news/include/inn/defines.h
news/include/inn/dispatch.h
news/include/inn/hashtab.h
news/include/inn/history.h
news/include/inn/innconf.h
news/include/inn/list.h
news/include/inn/md5.h
news/include/inn/messages.h
news/include/inn/mmap.h
news/include/inn/network.h
news/include/inn/nntp.h
news/include/inn/options.h
news/include/inn/qio.h
news/include/inn/sequence.h
news/include/inn/system.h
news/include/inn/timer.h
news/include/inn/tst.h
news/include/inn/vector.h
news/include/inn/version.h
news/include/inn/wire.h
news/include/inndcomm.h
news/include/libinn.h
news/include/nntp.h
news/include/ov.h
news/include/paths.h
news/include/storage.h
news/lib/innreport_inn.pm
news/lib/innshellvars
news/lib/innshellvars.pl
news/lib/innshellvars.tcl
news/lib/libinn.a
news/lib/libinnhist.a
news/lib/libstorage.a
%%PORTDOCS%%%%DOCSDIR%%/CONTRIBUTORS
%%PORTDOCS%%%%DOCSDIR%%/HACKING
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/LICENSE
%%PORTDOCS%%%%DOCSDIR%%/MANIFEST
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@unexec rmdir %D/news/tmp 2>/dev/null || true
@dirrm news/spool/overview
@dirrm news/spool/outgoing
@dirrm news/spool/innfeed
@dirrm news/spool/incoming/bad
@dirrm news/spool/incoming
@dirrm news/spool/articles
@dirrm news/spool/archive
@dirrm news/spool/tmp
@dirrm news/spool
@dirrm news/run
@dirrm news/lib
@dirrm news/include/inn
@dirrm news/include
@unexec rmdir %D/%%ETCDIR%% 2>/dev/null || true
@dirrm news/doc
@dirrm news/db
@dirrm news/bin/rnews.libexec
@dirrm news/bin/filter
@dirrm news/bin/control
@dirrm news/bin/auth/resolv
@dirrm news/bin/auth/passwd
@dirrm news/bin/auth
@dirrm news/bin
@unexec rmdir %D/news 2>/dev/null || echo "===>  If you don't plan to use inn later, you can safely remove %D/news. Deinstall process kept these files intact:" | fmt && if [ -d %D/news ] ; then  for i in $(find %D/news -type f ); do echo "    - $i"; done;fi
