# net/balance/files/openrc-balance.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf to enable balance:
#
#balance_hosts="host1"
#balance_host1_flags="-a"
#balance_host1_address="host1.external.example"
#balance_host1_ports="http 8180"
#balance_host1_targets="host1.internal.example"
#
# See balance(8) for flags

name=balance
command=%%PREFIX%%/bin/balance
extra_commands="status"

start()
{
	if [ -x "${command}" ]; then
		for host in ${balance_hosts}; do
			eval ports=\"\${balance_${host}_ports}\"
			eval flags=\"\${balance_${host}_flags}\"
			eval address=\"\${balance_${host}_address}\"
			eval targets=\"\${balance_${host}_targets}\"
			if [ "" != "${address}" ]; then
				flags="${flags} -b ${address}"
			fi
			for port in ${ports}; do
				"${command}" ${flags} ${port} ${targets}
			done
		done
	fi
}

stop()
{
	if [ -x "${command}" ]; then
		for host in ${balance_hosts}; do
			eval ports=\"\${balance_${host}_ports}\"
			eval address=\"\${balance_${host}_address}\"
			flags=""
			if [ "" != "${address}" ]; then
				flags="-b ${address}"
			else
				address='*'
			fi
			for port in ${ports}; do
				einfo "balance at ${address}:${port}"
				"${command}" ${flags} -c kill ${port}
			done
		done
	fi
}

status()
{
	if [ -x "${command}" ]; then
		for host in ${balance_hosts}; do
			eval ports=\"\${balance_${host}_ports}\"
			eval address=\"\${balance_${host}_address}\"
			flags=""
			if [ "" != "${address}" ]; then
				flags="-b ${address}"
			else
				address='*'
			fi
			for port in ${ports}; do
				einfo"balance at ${address}:${port}"
				"${command}" ${flags} -c show ${port}
			done
		done
	fi
}

depend() {
	keyword -shutdown
}
