# sysutils/openhpi/files/openhpid.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: openhpid
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following to /etc/rc.conf[.local] to enable this service
#
# openhpid_enable="YES"
# openhpid_flags="-c %%PREFIX%%/etc/openhpi/openhpi.conf"
#

. /etc/rc.subr

name="openhpid"
rcvar="openhpid_enable"
command="%%PREFIX%%/sbin/openhpid"

: ${openhpid_enable=NO} ${openhpid_flags=-c %%PREFIX%%/etc/openhpi/openhpi.conf}

load_rc_config ${name}
run_rc_command "$1"
