# net/apinger/pkg-plist
sbin/apinger
etc/apinger.conf.sample
etc/rc.d/apinger.sh.sample
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%@dirrm %%DOCSDIR%%
