# emulators/rtc/pkg-plist
etc/rc.d/rtc.sh
include/rtc.h
modules/rtc.ko
%%EXAMPLESDIR%%/test.c
@dirrm %%EXAMPLESDIR%%
@dirrmtry modules
@unexec rm -f /dev/rtc
