# emulators/rtc/pkg-plist
etc/rc.d/rtc.sh
include/rtc.h
modules/rtc.ko
%%EXAMPLESDIR%%/test.c
@postexec rm -f %%DEVDIR%%/rtc
@postunexec rm -f /dev/rtc
