# emulators/rtc/pkg-plist
etc/rc.d/rtc.sh
include/rtc.h
modules/rtc.ko
%%EXAMPLESDIR%%/test.c
@exec rm -f %%DEVDIR%%/rtc
@unexec rm -f /dev/rtc
