# emulators/rtc/pkg-plist
etc/rc.d/rtc.sh
include/rtc.h
modules/rtc.ko
%%EXAMPLESDIR%%/test.c
%%MAKEDEV%%@exec mknod /dev/rtc c %%CDEV_MAJOR%% 0
@exec rm -f %%DEVDIR%%/rtc
@exec /usr/sbin/kldxref %D/modules
@unexec /usr/sbin/kldxref %D/modules
@dirrm %%EXAMPLESDIR%%
@dirrmtry modules
@unexec rm -f /dev/rtc
