# www/varnish4/files/varnishncsa.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: varnishncsa
# REQUIRE: DAEMON varnishd
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable varnishncsa:
#
# varnishncsa_enable="YES"
#
# Configuration variables and their default values:
#
# varnishncsa_pidfile - full path to the PID file.
#	default: "/var/run/varnishncsa.pid"
#
# varnishncsa_file - full path to the log file.
#	default: "/var/log/varnishncsa.log"
#
# varnishncsa_flags - command line arguments.
#	default: "-D -P ${varnishncsa_pidfile} -a -w ${varnishncsa_file}${varnishncsa_logformat:+ -P \"$varnishncsa_logformat\""
#
# varnishncsa_logformat - log file format.
#       default: "" (uses varnishncsa's default format)
#       example: "%h %l %u %t %r %s %b %{Referer}i %{User-agent}i"
#
# Add the following line to /etc/newsyslog.conf to rotate the log file
# once a day:
#
# /var/log/varnishncsa.log 640 7 * @T00 JB /var/run/varnishncsa.pid
#
# See varnishncsa(1) for a detailed overview of command-line options.
#

. /etc/rc.subr

name=varnishncsa
rcvar=varnishncsa_enable

load_rc_config ${name}

: ${varnishncsa_enable:=NO}
: ${varnishncsa_pidfile=/var/run/${name}.pid}
: ${varnishncsa_file=/var/log/${name}.log}
: ${varnishncsa_flags="-P ${varnishncsa_pidfile} -D -a -w ${varnishncsa_file}${varnishncsa_logformat:+ -F \"$varnishncsa_logformat\"}"}
: ${varnishncsa_user=varnishlog}

command="%%PREFIX%%/bin/${name}"
pidfile=${varnishncsa_pidfile}
start_precmd=precmd

precmd()
{
	if [ ! -e ${pidfile} ]; then
		install -o ${varnishncsa_user} -g varnish /dev/null ${pidfile};
	fi

	if [ ! -e ${varnishncsa_file} ]; then
		install -o ${varnishncsa_user} -g varnish /dev/null ${varnishncsa_file};
	fi

	waitcount=0
	vsm=%%PREFIX%%/varnish/$(hostname)/_.vsm
	
        # wait no longer than 5s for  _.vsm file to arrive from varnishd
	while :
	do
		if [ -f ${vsm} ] ; then
			break
		else
			case "${waitcount}" in
				5)	echo "${name}: Failed to start. Missing ${vsm} ?"
					break
					;;
				*)	echo "${name}: waiting for ${vsm}"
					waitcount=$((waitcount+1))
					sleep 1
					;;
			esac
		fi
	done
}

run_rc_command "$1"
