# misc/upclient/pkg-plist
@comment $FreeBSD$
etc/rc.d/upclient.sh
@unexec if cmp -s %D/etc/upclient.conf %D/etc/upclient.conf.sample; then rm -f %D/etc/upclient.conf; fi
etc/upclient.conf.sample
@exec [ -f %D/etc/upclient.conf ] || cp -p %D/etc/upclient.conf.sample %D/etc/upclient.conf
sbin/clientchk
sbin/upchk
sbin/upclient
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/FAQ
%%PORTDOCS%%%%DOCSDIR%%/HISTORY
%%PORTDOCS%%%%DOCSDIR%%/HISTORY-BETA
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%%%DOCSDIR%%/UptimesProtocol5
share/locale/de/LC_MESSAGES/upclient.mo
share/locale/es/LC_MESSAGES/upclient.mo
share/locale/fr/LC_MESSAGES/upclient.mo
share/locale/it/LC_MESSAGES/upclient.mo
share/locale/nl/LC_MESSAGES/upclient.mo
%%PORTDOCS%%@dirrm %%DOCSDIR%%
