# devel/storm/files/storm-supervisor.in
#!/bin/sh
#
# Copyright (c) 2014, Radim Kolar
# All rights reserved.
#
# Redistribution and use in source and binary forms, with or without
# modification, are permitted provided that the following conditions are met:
#
#  * Redistributions of source code must retain the above copyright notice,
#    this list of conditions and the following disclaimer.
#  * Redistributions in binary form must reproduce the above copyright
#    notice, this list of conditions and the following disclaimer in the
#    documentation and/or other materials provided with the distribution.
#
# THIS SOFTWARE IS PROVIDED BY THE AUTHOR AND CONTRIBUTORS ``AS IS'' AND ANY
# EXPRESS OR IMPLIED WARRANTIES, INCLUDING, BUT NOT LIMITED TO, THE IMPLIED
# WARRANTIES OF MERCHANTABILITY AND FITNESS FOR A PARTICULAR PURPOSE ARE
# DISCLAIMED. IN NO EVENT SHALL THE AUTHOR OR CONTRIBUTORS BE LIABLE FOR ANY
# DIRECT, INDIRECT, INCIDENTAL, SPECIAL, EXEMPLARY, OR CONSEQUENTIAL DAMAGES
# (INCLUDING, BUT NOT LIMITED TO, PROCUREMENT OF SUBSTITUTE GOODS OR
# SERVICES; LOSS OF USE, DATA, OR PROFITS; OR BUSINESS INTERRUPTION) HOWEVER
# CAUSED AND ON ANY THEORY OF LIABILITY, WHETHER IN CONTRACT, STRICT
# LIABILITY, OR TORT (INCLUDING NEGLIGENCE OR OTHERWISE) ARISING IN ANY WAY
# OUT OF THE USE OF THIS SOFTWARE, EVEN IF ADVISED OF THE POSSIBILITY OF SUCH
# DAMAGE.

# PROVIDE: storm-supervisor
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# $FreeBSD$

. /etc/rc.subr

name="storm_supervisor"
rcvar=storm_supervisor_enable
load_rc_config $name

: ${storm_supervisor_enable:="NO"}
: ${storm_supervisor_wrkdir:="/var/db/storm-supervisor"}

required_files="%%DATADIR%%/bin/storm %%DATADIR%%/conf/storm.yaml"
pidfile=/var/run/storm-supervisor.pid
command=/usr/sbin/daemon
command_args="-f -r -P ${pidfile} %%DATADIR%%/bin/storm supervisor"
start_precmd="${name}_prestart"
#procname="%%LOCALBASE%%/openjdk7/bin/java"

export PATH=$PATH:%%LOCALBASE%%/bin
export JAVA_VENDOR=openjdk
export JAVA_VERSION=1.7

storm_supervisor_prestart() {
	mkdir -p  ${storm_supervisor_wrkdir}
	chmod 740 ${storm_supervisor_wrkdir}
	cd ${storm_supervisor_wrkdir}
	return 0
}

run_rc_command "$1"
