# security/ssh/pkg-plist
etc/rc.d/sshd.sh
bin/scp1
bin/scp
bin/ssh1
bin/ssh
bin/slogin
bin/ssh-add1
bin/ssh-add
bin/ssh-agent1
bin/ssh-agent
bin/ssh-askpass1
bin/ssh-askpass
bin/ssh-keygen1
bin/ssh-keygen
bin/make-ssh-known-hosts1
bin/make-ssh-known-hosts
etc/ssh_config
etc/sshd_config
man/man1/make-ssh-known-hosts.1.gz
man/man1/scp.1.gz
man/man1/ssh-add.1.gz
man/man1/ssh-agent.1.gz
man/man1/ssh-keygen.1.gz
man/man1/ssh.1.gz
man/man1/slogin.1.gz
man/man8/sshd.8.gz
sbin/sshd1
sbin/sshd
@exec if [ ! -f %D/etc/ssh_host_key ]; then echo "Generating a secret host key.." ; %D/bin/ssh-keygen -N "" -f %D/etc/ssh_host_key; fi
