# security/ssh/pkg-plist
bin/make-ssh-known-hosts
bin/make-ssh-known-hosts1
bin/scp
bin/scp1
bin/slogin
bin/ssh
bin/ssh-add
bin/ssh-add1
bin/ssh-agent
bin/ssh-agent1
bin/ssh-keygen
bin/ssh-keygen1
bin/ssh1
etc/rc.d/sshd.sh
etc/ssh_config
etc/sshd_config
sbin/sshd
sbin/sshd1
@exec if [ ! -f %D/etc/ssh_host_key ]; then echo "Generating a secret host key.." ; %D/bin/ssh-keygen1 -N "" -f %D/etc/ssh_host_key; fi
