# mail/mailgraph/pkg-plist
etc/rc.d/mailgraph.sh
www/cgi-bin/mailgraph.cgi
sbin/mailgraph.pl
@unexec rmdir %D/www/cgi-bin 2>/dev/null || true
