# mail/mailgraph/pkg-plist
etc/rc.d/mailgraph.sh.sample
www/cgi-bin/mailgraph.cgi
www/cgi-bin/mailgraph.pl
@unexec rmdir %D/www/cgi-bin 2>/dev/null || true
@unexec rmdir %D/www 2>/dev/null || true
