# net-im/mastodon/files/mastodon_workers.in
#! /bin/sh

# PROVIDE: mastodon_workers
# REQUIRE: LOGIN postgresql nginx redis
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable the Mastodon background
# workers.
#
#  mastodon_workers_enable="YES"

. /etc/rc.subr

name=mastodon_workers
rcvar="${name}_enable"

load_rc_config "$name"

: ${mastodon_workers_enable:="NO"}
: ${mastodon_workers_loglevel:="error"}
: ${mastodon_workers_dbpool:="5"}
: ${mastodon_workers_threads:="5"}

mastodon_workers_chdir="%%WWWDIR%%"
mastodon_workers_env="DB_POOL=\"${mastodon_workers_dbpool}\" \
RAILS_ENV=production \
RAILS_LOG_LEVEL=\"${mastodon_workers_loglevel}\""
mastodon_workers_user="mastodon"
pidfile="%%WWWDIR%%/tmp/${name}.pid"
procname="%%RUBY%%:"

start_precmd="rm -f %%WWWDIR%%/Gemfile.lock"
command="%%PREFIX%%/bin/sidekiq"
command_flags="-c ${mastodon_workers_threads} -q default -q mailers -q pull -q push"
command_args=">> /var/log/mastodon/${name##mastodon_}.log 2>&1 &"

run_rc_command "$1"
