# net-im/mastodon/files/mastodon_stream.in
#!/bin/sh

# $FreeBSD$

# PROVIDE: mastodon_stream
# REQUIRE: LOGIN postgresql nginx redis
# KEYWORD: shutdown

# Add the following line to /etc/rc.conf to enable the Mastodon streaming
# service.
#
# mastodon_stream_enable="YES"
#
# mastodon_stream_port:       default 4000
#
# The options below can be set in /etc/rc.conf.  If they are not set there, they
# will be read from .env.production.
#
# mastodon_db_name:           e.g., mastodon_production
# mastodon_db_host:           e.g., localhost
# mastodon_db_user:           e.g., mastodon
# mastodon_db_port:           e.g., 5432
# mastodon_local_domain:      e.g., mastodon.mydomain.org (once set, do not change)
# mastodon_smtp_from_address: e.g., mastodon@mydomain.org
# mastodon_redis_host:        e.g., localhost
# mastodon_redis_port:        e.g., 6397
#

. /etc/rc.subr

desc="Mastodon streaming server"
name=mastodon_stream
rcvar=mastodon_stream_enable

load_rc_config $name

: ${mastodon_stream_enable:=NO}
: ${mastodon_stream_port=4000}

sig_stop=HUP
mastodon_stream_chdir="%%WWWDIR%%"
mastodon_stream_env="HOME=%%WWWDIR%% \
NODE_ENV=production \
PORT=\"${mastodon_stream_port}\" \
${mastodon_db_url:+DATABASE_URL=\"${mastodon_db_url}\"} \
${mastodon_db_host:+DB_HOST=\"${mastodon_db_host}\"} \
${mastodon_db_name:+DB_NAME=\"${mastodon_db_name}\"} \
${mastodon_db_port:+DB_PORT=\"${mastodon_db_port}\"} \
${mastodon_db_user:+DB_USER=\"${mastodon_db_user}\"} \
${mastodon_db_password:+DB_PASS=\"${mastodon_db_password}\"} \
${mastodon_local_domain:+LOCAL_DOMAIN=\"${mastodon_local_domain}\"} \
${mastodon_smtp_from_address:+SMTP_FROM_ADDRESS=\"${mastodon_smtp_from_address}\"} \
${mastodon_redis_host:+REDIS_HOST=\"${mastodon_redis_host}\"} \
${mastodon_redis_port:+REDIS_PORT=\"${mastodon_redis_port}\"} \
USER=mastodon"
mastodon_stream_user=mastodon

command="%%PREFIX%%/bin/node"
command_args="./streaming/index.js \
>> /var/log/mastodon/${name##mastodon_}.log 2>&1 &"

run_rc_command "$1"
