# security/bitwarden_rs/pkg-plist
bin/bitwarden_rs
etc/rc.conf.d/bitwarden_rs.sample
etc/rc.d/bitwarden_rs
%%WWWDIR%%/web-vault/.nojekyll
%%WWWDIR%%/web-vault/404.css
%%WWWDIR%%/web-vault/404.html
%%WWWDIR%%/web-vault/app-id.json
%%WWWDIR%%/web-vault/app/main.635f774a807de51e2fd8.css
%%WWWDIR%%/web-vault/app/main.635f774a807de51e2fd8.css.map
%%WWWDIR%%/web-vault/app/main.635f774a807de51e2fd8.js
%%WWWDIR%%/web-vault/app/main.635f774a807de51e2fd8.js.map
%%WWWDIR%%/web-vault/app/polyfills.635f774a807de51e2fd8.js
%%WWWDIR%%/web-vault/app/polyfills.635f774a807de51e2fd8.js.map
%%WWWDIR%%/web-vault/app/vendor.635f774a807de51e2fd8.js
%%WWWDIR%%/web-vault/app/vendor.635f774a807de51e2fd8.js.map
%%WWWDIR%%/web-vault/browserconfig.xml
%%WWWDIR%%/web-vault/connectors/duo.635f774a807de51e2fd8.css
%%WWWDIR%%/web-vault/connectors/duo.635f774a807de51e2fd8.css.map
%%WWWDIR%%/web-vault/connectors/duo.635f774a807de51e2fd8.js
%%WWWDIR%%/web-vault/connectors/duo.635f774a807de51e2fd8.js.map
%%WWWDIR%%/web-vault/connectors/sso.635f774a807de51e2fd8.css
%%WWWDIR%%/web-vault/connectors/sso.635f774a807de51e2fd8.css.map
%%WWWDIR%%/web-vault/connectors/sso.635f774a807de51e2fd8.js
%%WWWDIR%%/web-vault/connectors/sso.635f774a807de51e2fd8.js.map
%%WWWDIR%%/web-vault/connectors/u2f.635f774a807de51e2fd8.js
%%WWWDIR%%/web-vault/connectors/u2f.635f774a807de51e2fd8.js.map
%%WWWDIR%%/web-vault/duo-connector.html
%%WWWDIR%%/web-vault/favicon.ico
%%WWWDIR%%/web-vault/fonts/Open_Sans-italic-300.woff
%%WWWDIR%%/web-vault/fonts/Open_Sans-italic-400.woff
%%WWWDIR%%/web-vault/fonts/Open_Sans-italic-600.woff
%%WWWDIR%%/web-vault/fonts/Open_Sans-italic-700.woff
%%WWWDIR%%/web-vault/fonts/Open_Sans-italic-800.woff
%%WWWDIR%%/web-vault/fonts/Open_Sans-normal-300.woff
%%WWWDIR%%/web-vault/fonts/Open_Sans-normal-400.woff
%%WWWDIR%%/web-vault/fonts/Open_Sans-normal-600.woff
%%WWWDIR%%/web-vault/fonts/Open_Sans-normal-700.woff
%%WWWDIR%%/web-vault/fonts/Open_Sans-normal-800.woff
%%WWWDIR%%/web-vault/fonts/fontawesome-webfont.eot
%%WWWDIR%%/web-vault/fonts/fontawesome-webfont.svg
%%WWWDIR%%/web-vault/fonts/fontawesome-webfont.ttf
%%WWWDIR%%/web-vault/fonts/fontawesome-webfont.woff
%%WWWDIR%%/web-vault/fonts/fontawesome-webfont.woff2
%%WWWDIR%%/web-vault/images/0.png
%%WWWDIR%%/web-vault/images/1.png
%%WWWDIR%%/web-vault/images/2.png
%%WWWDIR%%/web-vault/images/3.png
%%WWWDIR%%/web-vault/images/4.png
%%WWWDIR%%/web-vault/images/404.png
%%WWWDIR%%/web-vault/images/cards.png
%%WWWDIR%%/web-vault/images/fa-globe.png
%%WWWDIR%%/web-vault/images/icons/android-chrome-192x192.png
%%WWWDIR%%/web-vault/images/icons/android-chrome-512x512.png
%%WWWDIR%%/web-vault/images/icons/apple-touch-icon.png
%%WWWDIR%%/web-vault/images/icons/favicon-16x16.png
%%WWWDIR%%/web-vault/images/icons/favicon-32x32.png
%%WWWDIR%%/web-vault/images/icons/mstile-150x150.png
%%WWWDIR%%/web-vault/images/icons/safari-pinned-tab.svg
%%WWWDIR%%/web-vault/images/loading.svg
%%WWWDIR%%/web-vault/images/logo-dark@2x.png
%%WWWDIR%%/web-vault/images/logo-horizontal-white.png
%%WWWDIR%%/web-vault/images/register-layout/logo-horizontal-white.png
%%WWWDIR%%/web-vault/images/register-layout/wired-logo.png
%%WWWDIR%%/web-vault/images/totp-countdown.png
%%WWWDIR%%/web-vault/images/two-factor/0.png
%%WWWDIR%%/web-vault/images/two-factor/1.png
%%WWWDIR%%/web-vault/images/two-factor/2.png
%%WWWDIR%%/web-vault/images/two-factor/3.png
%%WWWDIR%%/web-vault/images/two-factor/4.png
%%WWWDIR%%/web-vault/images/two-factor/6.png
%%WWWDIR%%/web-vault/images/u2fkey.jpg
%%WWWDIR%%/web-vault/images/wired-logo.png
%%WWWDIR%%/web-vault/images/yubikey.jpg
%%WWWDIR%%/web-vault/index.html
%%WWWDIR%%/web-vault/locales/af/messages.json
%%WWWDIR%%/web-vault/locales/be/messages.json
%%WWWDIR%%/web-vault/locales/bg/messages.json
%%WWWDIR%%/web-vault/locales/ca/messages.json
%%WWWDIR%%/web-vault/locales/cs/messages.json
%%WWWDIR%%/web-vault/locales/da/messages.json
%%WWWDIR%%/web-vault/locales/de/messages.json
%%WWWDIR%%/web-vault/locales/el/messages.json
%%WWWDIR%%/web-vault/locales/en/messages.json
%%WWWDIR%%/web-vault/locales/en_GB/messages.json
%%WWWDIR%%/web-vault/locales/en_IN/messages.json
%%WWWDIR%%/web-vault/locales/eo/messages.json
%%WWWDIR%%/web-vault/locales/es/messages.json
%%WWWDIR%%/web-vault/locales/et/messages.json
%%WWWDIR%%/web-vault/locales/fi/messages.json
%%WWWDIR%%/web-vault/locales/fr/messages.json
%%WWWDIR%%/web-vault/locales/he/messages.json
%%WWWDIR%%/web-vault/locales/hr/messages.json
%%WWWDIR%%/web-vault/locales/hu/messages.json
%%WWWDIR%%/web-vault/locales/id/messages.json
%%WWWDIR%%/web-vault/locales/it/messages.json
%%WWWDIR%%/web-vault/locales/ja/messages.json
%%WWWDIR%%/web-vault/locales/ko/messages.json
%%WWWDIR%%/web-vault/locales/lv/messages.json
%%WWWDIR%%/web-vault/locales/ml/messages.json
%%WWWDIR%%/web-vault/locales/nb/messages.json
%%WWWDIR%%/web-vault/locales/nl/messages.json
%%WWWDIR%%/web-vault/locales/pl/messages.json
%%WWWDIR%%/web-vault/locales/pt_BR/messages.json
%%WWWDIR%%/web-vault/locales/pt_PT/messages.json
%%WWWDIR%%/web-vault/locales/ro/messages.json
%%WWWDIR%%/web-vault/locales/ru/messages.json
%%WWWDIR%%/web-vault/locales/si/messages.json
%%WWWDIR%%/web-vault/locales/sk/messages.json
%%WWWDIR%%/web-vault/locales/sl/messages.json
%%WWWDIR%%/web-vault/locales/sr/messages.json
%%WWWDIR%%/web-vault/locales/sv/messages.json
%%WWWDIR%%/web-vault/locales/tr/messages.json
%%WWWDIR%%/web-vault/locales/uk/messages.json
%%WWWDIR%%/web-vault/locales/vi/messages.json
%%WWWDIR%%/web-vault/locales/zh_CN/messages.json
%%WWWDIR%%/web-vault/locales/zh_TW/messages.json
%%WWWDIR%%/web-vault/manifest.json
%%WWWDIR%%/web-vault/scripts/dropin.js
%%WWWDIR%%/web-vault/scripts/qrious.min.js
%%WWWDIR%%/web-vault/scripts/u2f.js
%%WWWDIR%%/web-vault/sso-connector.html
%%WWWDIR%%/web-vault/u2f-connector.html
%%WWWDIR%%/web-vault/version.json
@dir %%ETCDIR%%/rc.conf.d
@dir %%WWWDIR%%/data
