# devel/linux-js/files/jscal.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: jscal
# REQUIRE: linux_js
# KEYWORD: shutdown

. /etc/rc.subr

name=jscal
rcvar=jscal_enable

load_rc_config $name

jscal_enable=${jscal_enable:-"NO"}
jscal_devices="${jscal_devices:-/dev/input/js*}"
jscal_state="${jscal_state:-/var/db/jscal-state}"

command="%%PREFIX%%/bin/jscal"
start_cmd="jscal_start"
stop_cmd="jscal_stop"

jscal_start()
{
	echo "Starting ${name}."
	for d in $jscal_devices; do
		[ "x$d" = "x$jscal_devices" ] && break
		statefile="$jscal_state-`basename $d`"
		[ -f "$statefile" ] && . "$statefile" 2>/dev/null
	done
}

jscal_stop()
{
	echo "Stopping ${name}."
	for d in $jscal_devices; do
		[ "x$d" = "x$jscal_devices" ] && break
		statefile="$jscal_state-`basename $d`"
		cal=`"$command" -p "$d" 2>/dev/null` && \
		echo "$cal" | sed -e "s|^jscal|$command|" > "$statefile"
	done
}

run_rc_command "$1"
