# comms/conserver/files/conserver.in
#!/bin/sh

# $FreeBSD$

# PROVIDE: conserver
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

# To enable conserver, add lines like below to /etc/rc.conf.
#
# conserver_enable="YES"

conserver_enable=${conserver_enable:-"no"}

. /etc/rc.subr

name=conserver
rcvar=conserver_enable

command=%%PREFIX%%/sbin/conserver
conserver_configfile=${conserver_configfile:-"%%PREFIX%%/etc/conserver.cf"}
command_args="-d -n -C ${conserver_configfile}"
required_files=${conserver_configfile}

load_rc_config $name

case "$1" in
start)	if [ -x ${command} -a -f ${configfile} ]; then
		run_rc_command "$1"
	fi
	;;
*)
	run_rc_command "$1"
	;;
esac

exit 0
