# databases/pgcluster/files/pgreplicate.in
#!/bin/sh
#
# $FreeBSD: /tmp/pcvs/ports/databases/pgcluster/files/Attic/pgreplicate.in,v 1.3 2010-03-27 00:12:48 dougb Exp $
#
# PROVIDE: pgreplicate
# REQUIRE: DAEMON
# BEFORE:  pglb
#
# Add the following line to /etc/rc.conf to enable pgreplicate:
#
# pgreplicate_enable="YES"
#

. /etc/rc.subr

name="pgreplicate"
rcvar=`set_rcvar`

command=%%PREFIX%%/bin/${name}

load_rc_config $name

: ${pgreplicate_enable="NO"}
: ${pgreplicate_data="%%PREFIX%%/pgsql/data"}

pgreplicate_user=pgsql
pgreplicate_flags="${pgreplicate_flags} -D ${pgreplicate_data}"
pidfile="${pgreplicate_data}/${name}.pid"

stop_cmd="pgreplicate_stop"

pgreplicate_stop()
{
    su -m ${pgreplicate_user} -c "exec ${command} -D ${pgreplicate_data} stop"
}

run_rc_command "$1"
