# net-mgmt/chillispot/pkg-plist
@unexec %D/etc/rc.d/chillispot forcestop 2>/dev/null || true
sbin/chilli
%%DATADIR%%/chilli.conf.sample
%%DATADIR%%/dictionary.chillispot
%%DATADIR%%/freeradius.users
%%DATADIR%%/hotspotlogin.cgi
%%DATADIR%%/pf.conf.sample
%%DATADIR%%/ipfw-config.sample
%%DATADIR%%/installguide.txt
@dirrm %%DATADIR%%
