# misc/proxyper/pkg-plist
etc/rc.d/proxyper.sh
proxyper/ChangeLog.txt
proxyper/Readme.txt
proxyper/proxyper
proxyper/proxyper.ini
proxyper/proxyper.sh
