# misc/proxyper/pkg-plist
etc/rc.d/proxyper.sh
distributed.net/proxy/ChangeLog.txt
distributed.net/proxy/Readme.txt
distributed.net/proxy/proxyper
distributed.net/proxy/proxyper.ini
distributed.net/proxy/proxyper.sh
