# misc/proxyper/pkg-plist
distributed.net/proxy/ChangeLog.txt
distributed.net/proxy/Readme.txt
distributed.net/proxy/proxyper
distributed.net/proxy/proxyper.ini
distributed.net/proxy/proxyper.sh
@exec ln -sf %D/distributed.net/proxy/proxyper.sh %D/etc/rc.d/proxyper.sh
@unexec rm %D/etc/rc.d/proxyper.sh
@dirrm distributed.net/proxy
