# misc/proxyper/pkg-plist
distributed.net/proxy/ChangeLog.txt
distributed.net/proxy/manual.html
distributed.net/proxy/proxyper
distributed.net/proxy/proxyper.ini
distributed.net/proxy/proxyper.sh
distributed.net/proxy/proxyper.sh.bak
@exec ln -sf %D/distributed.net/proxy/proxyper.sh %D/etc/rc.d/proxyper.sh
@unexec rm %D/etc/rc.d/proxyper.sh
@dirrm distributed.net/proxy
@dirrm distributed.net
