# sysutils/warden/pkg-plist
bin/warden
etc/warden.conf.dist
etc/rc.d/wardenrc
%%DATADIR%%/agent.png
%%DATADIR%%/warden.png
%%DATADIR%%/scripts/rc.d/wardenrc
%%DATADIR%%/scripts/rc.d/startup.sh
%%DATADIR%%/scripts/rc.d/shutdown.sh
%%DATADIR%%/scripts/backend/checkupdates.sh
%%DATADIR%%/scripts/backend/createjail.sh
%%DATADIR%%/scripts/backend/createtemplate.sh
%%DATADIR%%/scripts/backend/zfslistsnap.sh
%%DATADIR%%/scripts/backend/zfsrevertsnap.sh
%%DATADIR%%/scripts/backend/details.sh
%%DATADIR%%/scripts/backend/zfsmksnap.sh
%%DATADIR%%/scripts/backend/startjail.sh
%%DATADIR%%/scripts/backend/chrootjail.sh
%%DATADIR%%/scripts/backend/listjails.sh
%%DATADIR%%/scripts/backend/listpbis.sh
%%DATADIR%%/scripts/backend/cronsnap.sh
%%DATADIR%%/scripts/backend/setupusers.sh
%%DATADIR%%/scripts/backend/checkstatus.sh
%%DATADIR%%/scripts/backend/functions.sh
%%DATADIR%%/scripts/backend/fbsdupdate.sh
%%DATADIR%%/scripts/backend/fbsdupgrade.sh
%%DATADIR%%/scripts/backend/pkgupdate.sh
%%DATADIR%%/scripts/backend/importjail.sh
%%DATADIR%%/scripts/backend/stopjail.sh
%%DATADIR%%/scripts/backend/zfsrmsnap.sh
%%DATADIR%%/scripts/backend/zfsrmsnapdate.sh
%%DATADIR%%/scripts/backend/listpkgs.sh
%%DATADIR%%/scripts/backend/exportjail.sh
%%DATADIR%%/scripts/backend/deletejail.sh
%%DATADIR%%/scripts/backend/zfsclone.sh
%%DATADIR%%/scripts/backend/zfscronsnap.sh
%%DATADIR%%/linux-installs/debian-7-wheezy
%%DATADIR%%/linux-installs/gentoo-stage3-i486
