# sysutils/warden/pkg-plist
bin/warden
share/warden/agent.png
share/warden/warden.png
share/warden/scripts/rc.d/wardenrc
share/warden/scripts/rc.d/startup.sh
share/warden/scripts/rc.d/shutdown.sh
share/warden/scripts/backend/checkupdates.sh
share/warden/scripts/backend/createjail.sh
share/warden/scripts/backend/createtemplate.sh
share/warden/scripts/backend/zfslistsnap.sh
share/warden/scripts/backend/zfsrevertsnap.sh
share/warden/scripts/backend/details.sh
share/warden/scripts/backend/zfsmksnap.sh
share/warden/scripts/backend/startjail.sh
share/warden/scripts/backend/chrootjail.sh
share/warden/scripts/backend/listjails.sh
share/warden/scripts/backend/listpbis.sh
share/warden/scripts/backend/cronsnap.sh
share/warden/scripts/backend/setupusers.sh
share/warden/scripts/backend/checkstatus.sh
share/warden/scripts/backend/functions.sh
share/warden/scripts/backend/importjail.sh
share/warden/scripts/backend/stopjail.sh
share/warden/scripts/backend/zfsrmsnap.sh
share/warden/scripts/backend/zfsrmsnapdate.sh
share/warden/scripts/backend/listpkgs.sh
share/warden/scripts/backend/exportjail.sh
share/warden/scripts/backend/deletejail.sh
share/warden/scripts/backend/zfsclone.sh
share/warden/scripts/backend/zfscronsnap.sh
share/warden/linux-installs/debian-7-wheezy
share/warden/linux-installs/gentoo-stage3-i486
etc/warden.conf.dist
etc/rc.d/wardenrc
