# games/tshock/pkg-plist
%%DAEMON%%etc/rc.d/tshock
%%STANDALONE%%bin/tshock
%%DATADIR%%/bin/e_sqlite3.so
%%DATADIR%%/GeoIP.dat
%%DATADIR%%/ServerPlugins/TShockAPI.deps.json
%%DATADIR%%/ServerPlugins/TShockAPI.dll
%%DATADIR%%/ServerPlugins/TShockAPI.pdb
%%DATADIR%%/ServerPlugins/TShockAPI.xml
%%DATADIR%%/TShock.Server
%%DATADIR%%/TShock.Server.deps.json
%%DATADIR%%/TShock.Server.dll
%%DATADIR%%/TShock.Server.runtimeconfig.json
%%DATADIR%%/bin/BCrypt.Net-Next.dll
%%DATADIR%%/bin/BouncyCastle.Cryptography.dll
%%DATADIR%%/bin/GetText.dll
%%DATADIR%%/bin/Google.Protobuf.dll
%%DATADIR%%/bin/HttpServer.dll
%%DATADIR%%/bin/HttpServer.pdb
%%DATADIR%%/bin/HttpServer.xml
%%DATADIR%%/bin/K4os.Compression.LZ4.Streams.dll
%%DATADIR%%/bin/K4os.Compression.LZ4.dll
%%DATADIR%%/bin/K4os.Hash.xxHash.dll
%%DATADIR%%/bin/Microsoft.Data.Sqlite.dll
%%DATADIR%%/bin/Microsoft.Extensions.DependencyInjection.Abstractions.dll
%%DATADIR%%/bin/Microsoft.Extensions.Logging.Abstractions.dll
%%DATADIR%%/bin/ModFramework.dll
%%DATADIR%%/bin/Mono.Cecil.Mdb.dll
%%DATADIR%%/bin/Mono.Cecil.Pdb.dll
%%DATADIR%%/bin/Mono.Cecil.Rocks.dll
%%DATADIR%%/bin/Mono.Cecil.dll
%%DATADIR%%/bin/MonoMod.Backports.dll
%%DATADIR%%/bin/MonoMod.Core.dll
%%DATADIR%%/bin/MonoMod.ILHelpers.dll
%%DATADIR%%/bin/MonoMod.Iced.dll
%%DATADIR%%/bin/MonoMod.RuntimeDetour.dll
%%DATADIR%%/bin/MonoMod.Utils.dll
%%DATADIR%%/bin/MonoMod.dll
%%DATADIR%%/bin/MySql.Data.dll
%%DATADIR%%/bin/Newtonsoft.Json.dll
%%DATADIR%%/bin/Npgsql.dll
%%DATADIR%%/bin/NuGet.Common.dll
%%DATADIR%%/bin/NuGet.Configuration.dll
%%DATADIR%%/bin/NuGet.Frameworks.dll
%%DATADIR%%/bin/NuGet.Packaging.dll
%%DATADIR%%/bin/NuGet.Protocol.dll
%%DATADIR%%/bin/NuGet.Resolver.dll
%%DATADIR%%/bin/NuGet.Versioning.dll
%%DATADIR%%/bin/OTAPI.Runtime.dll
%%DATADIR%%/bin/OTAPI.dll
%%DATADIR%%/bin/SQLitePCLRaw.batteries_v2.dll
%%DATADIR%%/bin/SQLitePCLRaw.core.dll
%%DATADIR%%/bin/SQLitePCLRaw.provider.e_sqlite3.dll
%%DATADIR%%/bin/Steamworks.NET.dll
%%DATADIR%%/bin/System.CommandLine.dll
%%DATADIR%%/bin/System.Configuration.ConfigurationManager.dll
%%DATADIR%%/bin/System.Diagnostics.EventLog.dll
%%DATADIR%%/bin/System.Diagnostics.PerformanceCounter.dll
%%DATADIR%%/bin/System.IO.Packaging.dll
%%DATADIR%%/bin/System.Reflection.MetadataLoadContext.dll
%%DATADIR%%/bin/System.Security.Cryptography.Pkcs.dll
%%DATADIR%%/bin/System.Security.Cryptography.ProtectedData.dll
%%DATADIR%%/bin/System.Security.Permissions.dll
%%DATADIR%%/bin/System.Windows.Extensions.dll
%%DATADIR%%/bin/TShockPluginManager.dll
%%DATADIR%%/bin/TShockPluginManager.pdb
%%DATADIR%%/bin/TerrariaServer.dll
%%DATADIR%%/bin/TerrariaServer.pdb
%%DATADIR%%/bin/TerrariaServer.xml
%%DATADIR%%/bin/ZstdSharp.dll
%%DATADIR%%/cs/System.CommandLine.resources.dll
%%DATADIR%%/de/System.CommandLine.resources.dll
%%DATADIR%%/es/System.CommandLine.resources.dll
%%DATADIR%%/fr/System.CommandLine.resources.dll
%%DATADIR%%/i18n/de_DE/TShockAPI.mo
%%DATADIR%%/i18n/en_PT/TShockAPI.mo
%%DATADIR%%/i18n/es_ES/TShockAPI.mo
%%DATADIR%%/i18n/fr_FR/TShockAPI.mo
%%DATADIR%%/i18n/fr_QC/TShockAPI.mo
%%DATADIR%%/i18n/id_ID/TShockAPI.mo
%%DATADIR%%/i18n/it_IT/TShockAPI.mo
%%DATADIR%%/i18n/ja_JP/TShockAPI.mo
%%DATADIR%%/i18n/pl_PL/TShockAPI.mo
%%DATADIR%%/i18n/pt_BR/TShockAPI.mo
%%DATADIR%%/i18n/pt_PT/TShockAPI.mo
%%DATADIR%%/i18n/ru_RU/TShockAPI.mo
%%DATADIR%%/i18n/ru_UA/TShockAPI.mo
%%DATADIR%%/i18n/tok/TShockAPI.mo
%%DATADIR%%/i18n/tr_TR/TShockAPI.mo
%%DATADIR%%/i18n/zh_CN/TShockAPI.mo
%%DATADIR%%/it/System.CommandLine.resources.dll
%%DATADIR%%/ja/System.CommandLine.resources.dll
%%DATADIR%%/ko/System.CommandLine.resources.dll
%%DATADIR%%/pl/System.CommandLine.resources.dll
%%DATADIR%%/pt-BR/System.CommandLine.resources.dll
%%DATADIR%%/ru/System.CommandLine.resources.dll
%%DATADIR%%/tr/System.CommandLine.resources.dll
%%DATADIR%%/zh-Hans/System.CommandLine.resources.dll
%%DATADIR%%/zh-Hant/System.CommandLine.resources.dll
@dir %%ETCDIR%%
@dir /var/db/tshock
@dir /var/log/tshock
