# devel/py-mercurialserver/pkg-plist
etc/mercurialserver/access.conf
etc/mercurialserver/remote-hgrc.d/access.rc
etc/mercurialserver/remote-hgrc.d/logging.rc
@exec /bin/mkdir -p %D/etc/mercurialserver/keys/users
@exec /bin/mkdir -p %D/etc/mercurialserver/keys/root
