# sysutils/3dm/pkg-plist
@unexec %D/etc/rc.d/3dm2.sh stop > /dev/null 2>&1 || true
@unexec if cmp -s %D/etc/3dm2/3dm2.conf.sample %D/etc/3dm2/3dm2.conf; then rm -f %D/etc/3dm2/3dm2.conf; fi
etc/3dm2/3dm2.conf.sample
@exec if [ -f %B/3dm2.conf ] ; then cp -p %D/%F %B/3dm2.conf; fi
etc/3dm2/msg/tdm_msg_en
etc/3dm2/msg/tw_msg_en
sbin/3dm2
share/3dm2/en/3DM-Intro.4.10.html
share/3dm2/en/3DM-Intro.4.11.html
share/3dm2/en/3DM-Intro.4.12.html
share/3dm2/en/3DM-Intro.4.13.html
share/3dm2/en/3DM-Intro.4.14.html
share/3dm2/en/3DM-Intro.4.15.html
share/3dm2/en/3DM-Intro.4.16.html
share/3dm2/en/3DM-Intro.4.17.html
share/3dm2/en/3DM-Intro.4.18.html
share/3dm2/en/3DM-Intro.4.2.html
share/3dm2/en/3DM-Intro.4.3.html
share/3dm2/en/3DM-Intro.4.4.html
share/3dm2/en/3DM-Intro.4.5.html
share/3dm2/en/3DM-Intro.4.6.html
share/3dm2/en/3DM-Intro.4.7.html
share/3dm2/en/3DM-Intro.4.8.html
share/3dm2/en/3DM-Intro.4.9.html
share/3dm2/en/3DM2_Settings.html
share/3dm2/en/3DM_Ref.9.1.html
share/3dm2/en/3ware_GetStarted.3.1.html
share/3dm2/en/3ware_Intro.2.1.html
share/3dm2/en/3ware_Intro.2.2.html
share/3dm2/en/3ware_Intro.2.3.html
share/3dm2/en/3ware_Intro.2.4.html
share/3dm2/en/3ware_Intro.2.5.html
share/3dm2/en/3ware_Intro.2.6.html
share/3dm2/en/3ware_Intro.2.7.html
share/3dm2/en/3ware_Intro.2.8.html
share/3dm2/en/AEN_0001.html
share/3dm2/en/AEN_0002.html
share/3dm2/en/AEN_0003.html
share/3dm2/en/AEN_0004.html
share/3dm2/en/AEN_0005.html
share/3dm2/en/AEN_0006.html
share/3dm2/en/AEN_0007.html
share/3dm2/en/AEN_0008.html
share/3dm2/en/AEN_0009.html
share/3dm2/en/AEN_000A.html
share/3dm2/en/AEN_000B.html
share/3dm2/en/AEN_000C.html
share/3dm2/en/AEN_000D.html
share/3dm2/en/AEN_000E.html
share/3dm2/en/AEN_000F.html
share/3dm2/en/AEN_0019.html
share/3dm2/en/AEN_001A.html
share/3dm2/en/AEN_001E.html
share/3dm2/en/AEN_001F.html
share/3dm2/en/AEN_0021.html
share/3dm2/en/AEN_0022.html
share/3dm2/en/AEN_0023.html
share/3dm2/en/AEN_0024.html
share/3dm2/en/AEN_0025.html
share/3dm2/en/AEN_0026.html
share/3dm2/en/AEN_0027.html
share/3dm2/en/AEN_0028.html
share/3dm2/en/AEN_0029.html
share/3dm2/en/AEN_002A.html
share/3dm2/en/AEN_002B.html
share/3dm2/en/AEN_002C.html
share/3dm2/en/AEN_002D.html
share/3dm2/en/AEN_002E.html
share/3dm2/en/AEN_002F.html
share/3dm2/en/AEN_0030.html
share/3dm2/en/AEN_0032.html
share/3dm2/en/AEN_0033.html
share/3dm2/en/AEN_0034.html
share/3dm2/en/AEN_0035.html
share/3dm2/en/AEN_0036.html
share/3dm2/en/AEN_0037.html
share/3dm2/en/AEN_0038.html
share/3dm2/en/AEN_0039.html
share/3dm2/en/AEN_003A.html
share/3dm2/en/AEN_003B.html
share/3dm2/en/AEN_003C.html
share/3dm2/en/AEN_003D.html
share/3dm2/en/AEN_003E.html
share/3dm2/en/AEN_003F.html
share/3dm2/en/AEN_0040.html
share/3dm2/en/AEN_0041.html
share/3dm2/en/AEN_0042.html
share/3dm2/en/AEN_0043.html
share/3dm2/en/AEN_0044.html
share/3dm2/en/AEN_0045.html
share/3dm2/en/AEN_0046.html
share/3dm2/en/AEN_0047.html
share/3dm2/en/AEN_0048.html
share/3dm2/en/AEN_0049.html
share/3dm2/en/AEN_004A.html
share/3dm2/en/AEN_004B.html
share/3dm2/en/AEN_004C.html
share/3dm2/en/AEN_004D.html
share/3dm2/en/AEN_004E.html
share/3dm2/en/AEN_004F.html
share/3dm2/en/AEN_0050.html
share/3dm2/en/AEN_0051.html
share/3dm2/en/AEN_0052.html
share/3dm2/en/AEN_0053.html
share/3dm2/en/AEN_0055.html
share/3dm2/en/AEN_0056.html
share/3dm2/en/AEN_0057.html
share/3dm2/en/AEN_0058.html
share/3dm2/en/AEN_0059.html
share/3dm2/en/AEN_005A.html
share/3dm2/en/AEN_005B.html
share/3dm2/en/AEN_005C.html
share/3dm2/en/AEN_005D.html
share/3dm2/en/AEN_005E.html
share/3dm2/en/AEN_005F.html
share/3dm2/en/AMCC-logo-larger.jpg
share/3dm2/en/AMCC-logo-smaller.jpg
share/3dm2/en/AMCC-logo.jpg
share/3dm2/en/Alarms.html
share/3dm2/en/App_Glossary.11.1.html
share/3dm2/en/Battery_Backup.html
share/3dm2/en/Config_Controller.5.1.html
share/3dm2/en/Config_Controller.5.2.html
share/3dm2/en/Config_Controller.5.3.html
share/3dm2/en/Config_Controller.5.4.html
share/3dm2/en/Config_Controller.5.5.html
share/3dm2/en/Config_Controller.5.6.html
share/3dm2/en/Config_Controller.5.7.html
share/3dm2/en/Config_Controller.5.8.html
share/3dm2/en/Config_Controller.5.9.html
share/3dm2/en/Config_Units.6.1.html
share/3dm2/en/Config_Units.6.10.html
share/3dm2/en/Config_Units.6.11.html
share/3dm2/en/Config_Units.6.12.html
share/3dm2/en/Config_Units.6.13.html
share/3dm2/en/Config_Units.6.14.html
share/3dm2/en/Config_Units.6.15.html
share/3dm2/en/Config_Units.6.16.html
share/3dm2/en/Config_Units.6.17.html
share/3dm2/en/Config_Units.6.18.html
share/3dm2/en/Config_Units.6.19.html
share/3dm2/en/Config_Units.6.2.html
share/3dm2/en/Config_Units.6.20.html
share/3dm2/en/Config_Units.6.21.html
share/3dm2/en/Config_Units.6.22.html
share/3dm2/en/Config_Units.6.23.html
share/3dm2/en/Config_Units.6.24.html
share/3dm2/en/Config_Units.6.25.html
share/3dm2/en/Config_Units.6.26.html
share/3dm2/en/Config_Units.6.27.html
share/3dm2/en/Config_Units.6.3.html
share/3dm2/en/Config_Units.6.4.html
share/3dm2/en/Config_Units.6.5.html
share/3dm2/en/Config_Units.6.6.html
share/3dm2/en/Config_Units.6.7.html
share/3dm2/en/Config_Units.6.8.html
share/3dm2/en/Config_Units.6.9.html
share/3dm2/en/Controller_Details.html
share/3dm2/en/Controller_Settings.html
share/3dm2/en/Controller_Summary.html
share/3dm2/en/Copyright.12.1.html
share/3dm2/en/Drive_Details.html
share/3dm2/en/Drive_Information.html
share/3dm2/en/Enclosure_Details.html
share/3dm2/en/Enclosure_Information.html
share/3dm2/en/Exclamation.gif
share/3dm2/en/First.html
share/3dm2/en/Lightening.gif
share/3dm2/en/MaintainController.8.1.html
share/3dm2/en/MaintainController.8.10.html
share/3dm2/en/MaintainController.8.11.html
share/3dm2/en/MaintainController.8.12.html
share/3dm2/en/MaintainController.8.2.html
share/3dm2/en/MaintainController.8.3.html
share/3dm2/en/MaintainController.8.4.html
share/3dm2/en/MaintainController.8.5.html
share/3dm2/en/MaintainController.8.6.html
share/3dm2/en/MaintainController.8.7.html
share/3dm2/en/MaintainController.8.8.html
share/3dm2/en/MaintainController.8.9.html
share/3dm2/en/Maintain_Units.7.1.html
share/3dm2/en/Maintain_Units.7.10.html
share/3dm2/en/Maintain_Units.7.11.html
share/3dm2/en/Maintain_Units.7.12.html
share/3dm2/en/Maintain_Units.7.13.html
share/3dm2/en/Maintain_Units.7.14.html
share/3dm2/en/Maintain_Units.7.15.html
share/3dm2/en/Maintain_Units.7.16.html
share/3dm2/en/Maintain_Units.7.17.html
share/3dm2/en/Maintain_Units.7.18.html
share/3dm2/en/Maintain_Units.7.19.html
share/3dm2/en/Maintain_Units.7.2.html
share/3dm2/en/Maintain_Units.7.20.html
share/3dm2/en/Maintain_Units.7.21.html
share/3dm2/en/Maintain_Units.7.22.html
share/3dm2/en/Maintain_Units.7.23.html
share/3dm2/en/Maintain_Units.7.24.html
share/3dm2/en/Maintain_Units.7.25.html
share/3dm2/en/Maintain_Units.7.26.html
share/3dm2/en/Maintain_Units.7.27.html
share/3dm2/en/Maintain_Units.7.3.html
share/3dm2/en/Maintain_Units.7.4.html
share/3dm2/en/Maintain_Units.7.5.html
share/3dm2/en/Maintain_Units.7.6.html
share/3dm2/en/Maintain_Units.7.7.html
share/3dm2/en/Maintain_Units.7.8.html
share/3dm2/en/Maintain_Units.7.9.html
share/3dm2/en/Maintenance.html
share/3dm2/en/Note.gif
share/3dm2/en/OLH_3DM_2_First.html
share/3dm2/en/Scheduling.html
share/3dm2/en/Trbleshoot.10.1.html
share/3dm2/en/Trbleshoot.10.2.html
share/3dm2/en/Trbleshoot.10.3.html
share/3dm2/en/Trbleshoot.10.4.html
share/3dm2/en/Trbleshoot.10.5.html
share/3dm2/en/Unit_Details.html
share/3dm2/en/Unit_Information.html
share/3dm2/en/error_messages.html
share/3dm2/en/ix.html
share/3dm2/en/scripts/expand.js
share/3dm2/en/scripts/seealso.js
share/3dm2/en/sqbullet.gif
share/3dm2/en/toc.html
share/3dm2/en/css/3DM-Intro.css
share/3dm2/en/css/3DM_Ref.css
share/3dm2/en/css/3ware_GetStarted.css
share/3dm2/en/css/3ware_Intro.css
share/3dm2/en/css/AboutThisHelp.css
share/3dm2/en/css/App_Glossary.css
share/3dm2/en/css/Config_Controller.css
share/3dm2/en/css/Config_Units.css
share/3dm2/en/css/Copyright.css
share/3dm2/en/css/MaintainController.css
share/3dm2/en/css/Maintain_Units.css
share/3dm2/en/css/Trbleshoot.css
share/3dm2/en/css/webworks.css
share/3dm2/en/images/3DM-Intro.4.4.1.jpg
share/3dm2/en/images/3DM-Intro.4.8.1.jpg
share/3dm2/en/images/3DM-Settings_mac.JPG
share/3dm2/en/images/3DM_Ref.9.10.4.jpg
share/3dm2/en/images/3dm-Alarms-mac.jpg
share/3dm2/en/images/3dm-Controller_Details-mac.jpg
share/3dm2/en/images/3dm-Controller_Settings-mac.jpg
share/3dm2/en/images/3dm-Controller_Summary-mac-full.jpg
share/3dm2/en/images/3dm-Controller_Summary-mac.jpg
share/3dm2/en/images/3dm-Create-unit-1-mac.gif
share/3dm2/en/images/3dm-Drive_Information-mac.jpg
share/3dm2/en/images/3dm-Enclosure_Detail_Page.jpg
share/3dm2/en/images/3dm-Enclosure_Summary_Page-mac.jpg
share/3dm2/en/images/3dm-Schedule-drop-down-options-mac.gif
share/3dm2/en/images/3dm-Unit_Information-mac.jpg
share/3dm2/en/images/3dm-Unit_Information_Details-mac.jpg
share/3dm2/en/images/3dm-background-task-rate.jpg
share/3dm2/en/images/3dm-controller-drop-down.jpg
share/3dm2/en/images/3dm-create_unit_3.5mac.gif
share/3dm2/en/images/3dm-deletedrivesuc-mac.gif
share/3dm2/en/images/3dm-deletedunits-mac.gif
share/3dm2/en/images/3dm-drill-downs-pc.jpg
share/3dm2/en/images/3dm-drive-details-mac.JPG
share/3dm2/en/images/3dm-maint-page-mac-with-callouts-2-buttons.jpg
share/3dm2/en/images/3dm-maint-page-mac-with-callouts.jpg
share/3dm2/en/images/3dm-maintenance-PC-with-multiple-statuses.gif
share/3dm2/en/images/3dm-menu-mac.jpg
share/3dm2/en/images/3dm-scheduduling-mac.JPG
share/3dm2/en/images/3dm-schedule-select-type.gif
share/3dm2/en/images/3dm-selected-drives-to-create.gif
share/3dm2/en/images/3wUpdate-GetInfo.jpg
share/3dm2/en/images/3wUpdate-NoFrm1.jpg
share/3dm2/en/images/3wUpdate-Reboot.jpg
share/3dm2/en/images/3ware-update-final-new.JPG
share/3dm2/en/images/Config_Units.6.15.1.jpg
share/3dm2/en/images/Mac-OS-warning-before-parition.png
share/3dm2/en/images/MaintainController.8.12.2.jpg
share/3dm2/en/images/MaintainController.8.4.2.jpg
share/3dm2/en/images/MaintainController.8.4.3.jpg
share/3dm2/en/images/MaintainController.8.4.4.jpg
share/3dm2/en/images/MaintainController.8.4.5.jpg
share/3dm2/en/images/MaintainController.8.4.7.jpg
share/3dm2/en/images/MaintainController.8.5.1.jpg
share/3dm2/en/images/Maintain_Units.7.24.1.jpg
share/3dm2/en/images/Maintain_Units.7.25.1.jpg
share/3dm2/en/images/RAID-0.gif
share/3dm2/en/images/RAID-1.gif
share/3dm2/en/images/RAID-10.gif
share/3dm2/en/images/RAID-50-12-drives-select-drives-per-unit.gif
share/3dm2/en/images/RAID50Config.gif
share/3dm2/en/images/RAID5Config_Mac.gif
share/3dm2/en/images/RAID6Config.gif
share/3dm2/en/images/WinAVAlarm_config.gif
share/3dm2/en/images/WinAVAlarm_tray.gif
share/3dm2/en/images/XP-53.jpg
share/3dm2/en/images/XP-Update01.jpg
share/3dm2/en/images/XP-Update02.jpg
share/3dm2/en/images/XP-Update03.jpg
share/3dm2/en/images/XP-Update04.jpg
share/3dm2/en/images/XP-Update05.jpg
share/3dm2/en/images/XP-Update07.jpg
share/3dm2/en/images/XP-Update10.jpg
share/3dm2/en/images/XP-Update12.jpg
share/3dm2/en/images/alarms-sample-3dm.gif
share/3dm2/en/images/bbu.gif
share/3dm2/en/images/collapse.gif
share/3dm2/en/images/controller-settings.gif
share/3dm2/en/images/enable-write-cache.gif
share/3dm2/en/images/expanded.gif
share/3dm2/en/images/index.gif
share/3dm2/en/images/indexx.gif
share/3dm2/en/images/next.gif
share/3dm2/en/images/nextx.gif
share/3dm2/en/images/on-off-schedule-buttons.gif
share/3dm2/en/images/prev.gif
share/3dm2/en/images/prevx.gif
share/3dm2/en/images/seebttm.gif
share/3dm2/en/images/seelflow.gif
share/3dm2/en/images/seeright.gif
share/3dm2/en/images/seertlow.gif
share/3dm2/en/images/seertup.gif
share/3dm2/en/images/select-drive-during-rebuild-3dm.gif
share/3dm2/en/images/self-test-choices.gif
share/3dm2/en/images/toc.gif
share/3dm2/en/images/tocx.gif
@unexec unlink /etc/3dm2/3dm2.conf
@unexec unlink /etc/3dm2/3dm2.pem
@unexec unlink /etc/3dm2/msg
@dirrm share/3dm2/en/css
@dirrm share/3dm2/en/images
@dirrm share/3dm2/en/scripts
@dirrm share/3dm2/en
@dirrm share/3dm2
@unexec rmdir /etc/3dm2 2>/dev/null || true
@dirrmtry etc/3dm2/msg
@dirrmtry etc/3dm2
