# sysutils/3dm/pkg-plist
etc/rc.d/3dm2.sh
@unexec if cmp -s %D/etc/3dm2.conf.sample %D/etc/3dm2.conf; then rm -f %D/etc/3dm2.conf; fi
etc/3dm2.conf.sample
@exec [ -f %B/3dm2.conf ] || cp %B/%f %B/3dm2.conf
sbin/3dm2
share/3dm2/en/3DM_2_OLH-1-01.html
share/3dm2/en/3DM_2_OLH-1-02.html
share/3dm2/en/3DM_2_OLH-1-03.html
share/3dm2/en/3DM_2_OLH-1-04.html
share/3dm2/en/3DM_2_OLH-1-05.html
share/3dm2/en/3DM_2_OLH-1-06.html
share/3dm2/en/3DM_2_OLH-1-07.html
share/3dm2/en/3DM_2_OLH-1-08.html
share/3dm2/en/3DM_2_OLH-1-09.html
share/3dm2/en/3DM_2_OLH-1-10.html
share/3dm2/en/3DM_2_OLH-1-11.html
share/3dm2/en/3DM_2_OLH-1-12.html
share/3dm2/en/3DM_2_OLH-1-13.html
share/3dm2/en/3DM_2_OLH-1-14.html
share/3dm2/en/3DM_2_OLH-1-15.html
share/3dm2/en/3DM_2_OLH-1-16.html
share/3dm2/en/3DM_2_OLH-1-17.html
share/3dm2/en/3DM_2_OLH-1-18.html
share/3dm2/en/3DM_2_OLH-1-19.html
share/3dm2/en/3DM_2_OLH-1-20.html
share/3dm2/en/3DM_2_OLH-1-21.html
share/3dm2/en/3DM_2_OLH-1-22.html
share/3dm2/en/3DM_2_OLH-1-23.html
share/3dm2/en/3DM_2_OLH-1-24.html
share/3dm2/en/3DM_2_OLH-1-25.html
share/3dm2/en/3DM_2_OLH-1-26.html
share/3dm2/en/3DM_2_OLH-1-27.html
share/3dm2/en/3DM_2_OLH-1-28.html
share/3dm2/en/3DM_2_OLH-1-29.html
share/3dm2/en/3DM_2_OLH-1-30.html
share/3dm2/en/3DM_2_OLH-1-31.html
share/3dm2/en/3DM_2_OLH-2-1.html
share/3dm2/en/3DM_2_OLH-2-2.html
share/3dm2/en/3DM_2_OLH-2-3.html
share/3dm2/en/3DM_2_OLH-3-01.html
share/3dm2/en/3DM_2_OLH-3-02.html
share/3dm2/en/3DM_2_OLH-3-03.html
share/3dm2/en/3DM_2_OLH-3-04.html
share/3dm2/en/3DM_2_OLH-3-05.html
share/3dm2/en/3DM_2_OLH-3-06.html
share/3dm2/en/3DM_2_OLH-3-07.html
share/3dm2/en/3DM_2_OLH-3-08.html
share/3dm2/en/3DM_2_OLH-3-09.html
share/3dm2/en/3DM_2_OLH-3-10.html
share/3dm2/en/3DM_2_OLH-3-11.html
share/3dm2/en/3DM_2_OLH-3-12.html
share/3dm2/en/3DM_2_OLH-3-13.html
share/3dm2/en/3DM_2_OLH-3-14.html
share/3dm2/en/3DM_2_OLH-4-01.html
share/3dm2/en/3DM_2_OLH-4-02.html
share/3dm2/en/3DM_2_OLH-4-03.html
share/3dm2/en/3DM_2_OLH-4-04.html
share/3dm2/en/3DM_2_OLH-4-05.html
share/3dm2/en/3DM_2_OLH-4-06.html
share/3dm2/en/3DM_2_OLH-4-07.html
share/3dm2/en/3DM_2_OLH-4-08.html
share/3dm2/en/3DM_2_OLH-4-09.html
share/3dm2/en/3DM_2_OLH-4-10.html
share/3dm2/en/3DM_2_OLH-4-11.html
share/3dm2/en/3DM_2_OLH-4-12.html
share/3dm2/en/3DM_2_OLH-4-13.html
share/3dm2/en/3DM_2_OLH-4-14.html
share/3dm2/en/3DM_2_OLH-4-15.html
share/3dm2/en/3DM_2_OLH-4-16.html
share/3dm2/en/3DM_2_OLH-4-17.html
share/3dm2/en/3DM_2_OLH-4-18.html
share/3dm2/en/3DM_2_OLH-4-19.html
share/3dm2/en/3DM_2_OLH-4-20.html
share/3dm2/en/3DM_2_OLH-4-21.html
share/3dm2/en/3DM_2_OLH-4-22.html
share/3dm2/en/3DM_2_OLH-5-1.html
share/3dm2/en/3DM_2_OLH-5-2.html
share/3dm2/en/3DM_2_OLH-5-3.html
share/3dm2/en/3DM_2_OLH-5-4.html
share/3dm2/en/3DM_2_OLH-5-5.html
share/3dm2/en/3DM_2_OLH-6-1.html
share/3dm2/en/IX.html
share/3dm2/en/TOC.html
share/3dm2/en/catalog.css
share/3dm2/en/document.css
share/3dm2/en/images/3DM_2_OLH-3-08-2.gif
share/3dm2/en/images/3DM_2_OLH-3-09-3.gif
share/3dm2/en/images/3DM_2_OLH-3-10-4.gif
share/3dm2/en/images/3DM_2_OLH-3-12-5.gif
share/3dm2/en/images/3DM_2_OLH-3-13-6.gif
share/3dm2/en/images/3dm-selected-drives-to-create.gif
share/3dm2/en/images/3ware-logo.png
share/3dm2/en/images/3wblu_small50.gif
share/3dm2/en/images/3wblu_small75.gif
share/3dm2/en/images/AMCC-logo.GIF
share/3dm2/en/images/AMCC-logo.jpg
share/3dm2/en/images/AMCC.png
share/3dm2/en/images/RAID-50-12-drives-select-drives-per-unit.gif
share/3dm2/en/images/amcc-logo75.gif
share/3dm2/en/images/collapse.gif
share/3dm2/en/images/enable-write-cache.gif
share/3dm2/en/images/expanded.gif
share/3dm2/en/images/index.gif
share/3dm2/en/images/indexx.gif
share/3dm2/en/images/info-drill-downs.gif
share/3dm2/en/images/logo.gif
share/3dm2/en/images/menu-select.gif
share/3dm2/en/images/next.gif
share/3dm2/en/images/nextx.gif
share/3dm2/en/images/on-off-schedule-buttons.gif
share/3dm2/en/images/prev.gif
share/3dm2/en/images/prevx.gif
share/3dm2/en/images/seebttm.gif
share/3dm2/en/images/seelflow.gif
share/3dm2/en/images/seeright.gif
share/3dm2/en/images/seertlow.gif
share/3dm2/en/images/seertup.gif
share/3dm2/en/images/select-drive-during-rebuild-3dm.gif
share/3dm2/en/images/self-test-choices.gif
share/3dm2/en/images/summary-screen-with-callouts.gif
share/3dm2/en/images/toc.gif
share/3dm2/en/images/tocx.gif
share/3dm2/en/scripts/expand.js
share/3dm2/en/scripts/seealso.js
@unexec unlink /etc/3dm2/3dm2.conf
@unexec unlink /etc/3dm2/3dm2.pem
@unexec rmdir %D/share/3dm2/en/scripts %D/share/3dm2/en/images %D/share/3dm2/en %D/share/3dm2 /etc/3dm2 || true
