# sysutils/3dm/pkg-plist
etc/rc.d/3dmd.sh
etc/3dmd.conf
sbin/3dmd
share/3dm/3DM_browser.html
share/3dm/3DM_contact.html
share/3dm/3DM_email.html
share/3dm/3DM_help.html
share/3dm/3DM_hot_spare.html
share/3dm/3DM_hot_swap.html
share/3dm/3DM_linux_install.html
share/3dm/3DM_popup_linux.html
share/3dm/3DM_popup_windows.html
share/3dm/3DM_rebuild.html
share/3dm/3DM_windows_install.html
share/3dm/images/3DM-logo-no-transparency.gif
share/3dm/images/3DM-logo-transparency.gif
share/3dm/images/CAS_alarm.gif
share/3dm/images/CAS_details.gif
share/3dm/images/CAS_home.gif
share/3dm/images/CAS_monitor.gif
share/3dm/images/CAS_settings.gif
share/3dm/images/MAM_configure.gif
share/3dm/images/MAM_rebuild.gif
share/3dm/images/PNL_event_notification.gif
share/3dm/images/PNW_event_notification.gif
share/3dm/images/PNW_icon.gif
share/3dm/images/PNW_main.gif
share/3dm/images/PNW_snapshot.gif
share/3dm/images/WI_email_notification.gif
share/3dm/images/WI_remote_monitoring.gif
share/3dm/styles/3DMstyle.css
@unexec rmdir %D/share/3dm/styles %D/share/3dm/images %D/share/3dm || true
@unexec unlink /etc/3dmd.conf
