# sysutils/3dm/pkg-plist
@unexec %D/etc/rc.d/3dm2.sh stop > /dev/null 2>&1 || true
@unexec if cmp -s %D/etc/3dm2/3dm2.conf.sample %D/etc/3dm2/3dm2.conf; then rm -f %D/etc/3dm2/3dm2.conf; fi
etc/3dm2/3dm2.conf.sample
@exec if [ -f %B/3dm2.conf ] ; then cp -p %D/%F %B/3dm2.conf; fi
etc/3dm2/msg/tdm_msg_en
etc/3dm2/msg/tw_msg_en
sbin/3dm2
share/3dm2/en/3DM2_Settings.html
share/3dm2/en/3DM_2_OLH-1-01.html
share/3dm2/en/3DM_2_OLH-2-02.html
share/3dm2/en/3DM_2_OLH-2-03.html
share/3dm2/en/3DM_2_OLH-2-04.html
share/3dm2/en/3DM_2_OLH-2-05.html
share/3dm2/en/3DM_2_OLH-2-06.html
share/3dm2/en/3DM_2_OLH-2-07.html
share/3dm2/en/3DM_2_OLH-2-08.html
share/3dm2/en/3DM_2_OLH-2-09.html
share/3dm2/en/3DM_2_OLH-2-10.html
share/3dm2/en/3DM_2_OLH-2-11.html
share/3dm2/en/3DM_2_OLH-2-12.html
share/3dm2/en/3DM_2_OLH-2-13.html
share/3dm2/en/3DM_2_OLH-2-14.html
share/3dm2/en/3DM_2_OLH-2-15.html
share/3dm2/en/3DM_2_OLH-2-16.html
share/3dm2/en/3DM_2_OLH-2-17.html
share/3dm2/en/3DM_2_OLH-2-18.html
share/3dm2/en/3DM_2_OLH-2-19.html
share/3dm2/en/3DM_2_OLH-3-01.html
share/3dm2/en/3DM_2_OLH-3-02.html
share/3dm2/en/3DM_2_OLH-3-03.html
share/3dm2/en/3DM_2_OLH-3-04.html
share/3dm2/en/3DM_2_OLH-3-05.html
share/3dm2/en/3DM_2_OLH-3-06.html
share/3dm2/en/3DM_2_OLH-3-07.html
share/3dm2/en/3DM_2_OLH-3-08.html
share/3dm2/en/3DM_2_OLH-3-09.html
share/3dm2/en/3DM_2_OLH-3-10.html
share/3dm2/en/3DM_2_OLH-3-11.html
share/3dm2/en/3DM_2_OLH-4-1.html
share/3dm2/en/3DM_2_OLH-4-2.html
share/3dm2/en/3DM_2_OLH-4-3.html
share/3dm2/en/3DM_2_OLH-4-4.html
share/3dm2/en/3DM_2_OLH-4-5.html
share/3dm2/en/3DM_2_OLH-4-6.html
share/3dm2/en/3DM_2_OLH-4-7.html
share/3dm2/en/3DM_2_OLH-5-01.html
share/3dm2/en/3DM_2_OLH-5-02.html
share/3dm2/en/3DM_2_OLH-5-03.html
share/3dm2/en/3DM_2_OLH-5-04.html
share/3dm2/en/3DM_2_OLH-5-05.html
share/3dm2/en/3DM_2_OLH-5-06.html
share/3dm2/en/3DM_2_OLH-5-07.html
share/3dm2/en/3DM_2_OLH-5-08.html
share/3dm2/en/3DM_2_OLH-5-09.html
share/3dm2/en/3DM_2_OLH-5-10.html
share/3dm2/en/3DM_2_OLH-5-11.html
share/3dm2/en/3DM_2_OLH-5-12.html
share/3dm2/en/3DM_2_OLH-5-13.html
share/3dm2/en/3DM_2_OLH-5-14.html
share/3dm2/en/3DM_2_OLH-5-15.html
share/3dm2/en/3DM_2_OLH-5-16.html
share/3dm2/en/3DM_2_OLH-5-17.html
share/3dm2/en/3DM_2_OLH-5-18.html
share/3dm2/en/3DM_2_OLH-5-19.html
share/3dm2/en/3DM_2_OLH-5-20.html
share/3dm2/en/3DM_2_OLH-5-21.html
share/3dm2/en/3DM_2_OLH-5-22.html
share/3dm2/en/3DM_2_OLH-5-23.html
share/3dm2/en/3DM_2_OLH-6-01.html
share/3dm2/en/3DM_2_OLH-6-02.html
share/3dm2/en/3DM_2_OLH-6-03.html
share/3dm2/en/3DM_2_OLH-6-04.html
share/3dm2/en/3DM_2_OLH-6-05.html
share/3dm2/en/3DM_2_OLH-6-06.html
share/3dm2/en/3DM_2_OLH-6-07.html
share/3dm2/en/3DM_2_OLH-6-08.html
share/3dm2/en/3DM_2_OLH-6-09.html
share/3dm2/en/3DM_2_OLH-6-10.html
share/3dm2/en/3DM_2_OLH-6-11.html
share/3dm2/en/3DM_2_OLH-6-12.html
share/3dm2/en/3DM_2_OLH-6-13.html
share/3dm2/en/3DM_2_OLH-6-14.html
share/3dm2/en/3DM_2_OLH-6-15.html
share/3dm2/en/3DM_2_OLH-6-16.html
share/3dm2/en/3DM_2_OLH-6-17.html
share/3dm2/en/3DM_2_OLH-6-18.html
share/3dm2/en/3DM_2_OLH-6-19.html
share/3dm2/en/3DM_2_OLH-6-20.html
share/3dm2/en/3DM_2_OLH-6-21.html
share/3dm2/en/3DM_2_OLH-6-22.html
share/3dm2/en/3DM_2_OLH-6-23.html
share/3dm2/en/3DM_2_OLH-6-24.html
share/3dm2/en/3DM_2_OLH-7-1.html
share/3dm2/en/3DM_2_OLH-7-2.html
share/3dm2/en/3DM_2_OLH-7-3.html
share/3dm2/en/3DM_2_OLH-7-4.html
share/3dm2/en/3DM_2_OLH-7-5.html
share/3dm2/en/3DM_2_OLH-8-1.html
share/3dm2/en/3DM_2_OLH-8-2.html
share/3dm2/en/3DM_2_OLH-8-3.html
share/3dm2/en/3DM_2_OLH-8-4.html
share/3dm2/en/3DM_2_OLH-8-5.html
share/3dm2/en/3DM_2_OLH-8-6.html
share/3dm2/en/3DM_2_OLH-9-1.html
share/3dm2/en/Alarms.html
share/3dm2/en/Battery_Backup.html
share/3dm2/en/Controller_Details.html
share/3dm2/en/Controller_Settings.html
share/3dm2/en/Controller_Summary.html
share/3dm2/en/Drive_Details.html
share/3dm2/en/Drive_Information.html
share/3dm2/en/First.html
share/3dm2/en/IX.html
share/3dm2/en/Maintenance.html
share/3dm2/en/Scheduling.html
share/3dm2/en/TOC.html
share/3dm2/en/Unit_Details.html
share/3dm2/en/Unit_Information.html
share/3dm2/en/catalog.css
share/3dm2/en/document.css
share/3dm2/en/images/3DM_2_OLH-3-05-3.gif
share/3dm2/en/images/3DM_2_OLH-3-05-4.gif
share/3dm2/en/images/3dm-selected-drives-to-create.gif
share/3dm2/en/images/3dm-selected-drives-to-createkl.gif
share/3dm2/en/images/3ware-logo.png
share/3dm2/en/images/3wblu_small50.gif
share/3dm2/en/images/3wblu_small75.gif
share/3dm2/en/images/Migrate-2.gif
share/3dm2/en/images/Migrate0-5.gif
share/3dm2/en/images/RAID-0.gif
share/3dm2/en/images/RAID-1.gif
share/3dm2/en/images/RAID-50-12-drives-select-drives-per-unit.gif
share/3dm2/en/images/RAID-50-12-drives-select-drives-per-unitkm.gif
share/3dm2/en/images/alarms-sample-3dm.gif
share/3dm2/en/images/amcc-logo.gif
share/3dm2/en/images/amcc-logo.jpg
share/3dm2/en/images/amcc-logo75.gif
share/3dm2/en/images/amcc.png
share/3dm2/en/images/collapse.gif
share/3dm2/en/images/enable-write-cache.gif
share/3dm2/en/images/expanded.gif
share/3dm2/en/images/index.gif
share/3dm2/en/images/indexx.gif
share/3dm2/en/images/info-drill-downs.gif
share/3dm2/en/images/logo.gif
share/3dm2/en/images/menu-select.gif
share/3dm2/en/images/next.gif
share/3dm2/en/images/nextx.gif
share/3dm2/en/images/on-off-schedule-buttons.gif
share/3dm2/en/images/prev.gif
share/3dm2/en/images/prevx.gif
share/3dm2/en/images/seebttm.gif
share/3dm2/en/images/seelflow.gif
share/3dm2/en/images/seeright.gif
share/3dm2/en/images/seertlow.gif
share/3dm2/en/images/seertup.gif
share/3dm2/en/images/select-drive-during-rebuild-3dm.gif
share/3dm2/en/images/self-test-choices.gif
share/3dm2/en/images/summary-screen-with-callouts.gif
share/3dm2/en/images/toc.gif
share/3dm2/en/images/tocx.gif
@unexec unlink /etc/3dm2/3dm2.conf
@unexec unlink /etc/3dm2/3dm2.pem
@unexec unlink /etc/3dm2/msg
@dirrmtry share/3dm2/en/scripts
@dirrmtry share/3dm2/en/images
@dirrmtry share/3dm2/en
@dirrmtry share/3dm2
@unexec rmdir /etc/3dm2 2>/dev/null || true
@dirrmtry etc/3dm2/msg
@dirrmtry etc/3dm2
