# net/rwhoisd/pkg-plist
bin/rwhois
etc/rc.d/rwhoisd.sh
lib/rwhois/attribute_defs/asn.tmpl
lib/rwhois/attribute_defs/domain.tmpl
lib/rwhois/attribute_defs/host.tmpl
lib/rwhois/attribute_defs/network.tmpl
lib/rwhois/attribute_defs/nic.tmpl
lib/rwhois/attribute_defs/org.tmpl
lib/rwhois/attribute_defs/referral.tmpl
lib/rwhois/attribute_defs/soa.tmpl
lib/rwhois/attribute_defs/user.tmpl
lib/rwhois/bin/date
lib/rwhois/bin/rmkdbindex
lib/rwhois/bin/rwhois
lib/rwhois/bin/rwhoisparse
lib/rwhois/bin/sh
lib/rwhois/bin/sort
lib/rwhois/data/asn/asn.txt
lib/rwhois/data/asn/local-0.ndx
lib/rwhois/data/asn/local.db
lib/rwhois/data/asn/null.ca
lib/rwhois/data/domain/local-0.ndx
lib/rwhois/data/domain/local.db
lib/rwhois/data/domain/null.ca
lib/rwhois/data/domain/shoe-me.com.txt
lib/rwhois/data/host/local-0.ndx
lib/rwhois/data/host/local.db
lib/rwhois/data/host/null.ca
lib/rwhois/data/host/shoe_me.txt
lib/rwhois/data/network/local-0.ndx
lib/rwhois/data/network/local.db
lib/rwhois/data/network/null.ca
lib/rwhois/data/network/shoe-me.net.txt
lib/rwhois/data/referral/local-0.ndx
lib/rwhois/data/referral/local.db
lib/rwhois/data/referral/shoe-me-more.txt
lib/rwhois/data/soa/internic.txt
lib/rwhois/data/soa/local-0.ndx
lib/rwhois/data/soa/local.db
lib/rwhois/data/soa/null.ca
lib/rwhois/data/soa/shoe-me.txt
lib/rwhois/data/user/local-0.ndx
lib/rwhois/data/user/local.db
lib/rwhois/data/user/null.ca
lib/rwhois/data/user/shoe-me.be.txt
lib/rwhois/data/user/shoe-me.multiple.txt
lib/rwhois/data/user/shoe.gif
lib/rwhois/etc/resolv.conf
lib/rwhois/index_scripts/index_all
lib/rwhois/index_scripts/index_asn
lib/rwhois/index_scripts/index_domain
lib/rwhois/index_scripts/index_host
lib/rwhois/index_scripts/index_network
lib/rwhois/index_scripts/index_nic
lib/rwhois/index_scripts/index_org
lib/rwhois/index_scripts/index_referral
lib/rwhois/index_scripts/index_soa
lib/rwhois/index_scripts/index_user
lib/rwhois/output_templates/asn.full.output
lib/rwhois/output_templates/domain.full.output
lib/rwhois/output_templates/domain.sum.output
lib/rwhois/output_templates/host.contact.output
lib/rwhois/output_templates/host.full.output
lib/rwhois/output_templates/host.sum.output
lib/rwhois/output_templates/network.full.output
lib/rwhois/output_templates/network.ripe.output
lib/rwhois/output_templates/network.sum.output
lib/rwhois/output_templates/network.swip.output
lib/rwhois/output_templates/nic.full.output
lib/rwhois/output_templates/org.full.output
lib/rwhois/output_templates/referral.full.output
lib/rwhois/output_templates/soa.full.output
lib/rwhois/output_templates/user.contact.output
lib/rwhois/output_templates/user.full.output
lib/rwhois/output_templates/user.ripe.output
lib/rwhois/output_templates/user.sum.output
lib/rwhois/output_templates/user.swip.output
lib/rwhois/rwhois.allow
lib/rwhois/rwhois.conf
lib/rwhois/rwhois.conf.orig
lib/rwhois/rwhois.define
lib/rwhois/rwhois.deny
lib/rwhois/rwhois.directive
lib/rwhois/rwhois.display
lib/rwhois/rwhois.help
lib/rwhois/rwhois.log
lib/rwhois/rwhois.objects
lib/rwhois/rwhois.root
lib/rwhois/sbin/rwhoisd
lib/rwhois/usr/lib/libc.so.%%LIBCVERSION%%
lib/rwhois/usr/libexec/ld.so
share/doc/rwhois/EXAMPLES
share/doc/rwhois/INSTALLATION
share/doc/rwhois/NEWS
share/doc/rwhois/PROBLEMS
share/doc/rwhois/README
share/doc/rwhois/TODO
share/doc/rwhois/operations_guide.txt
share/doc/rwhois/sample.data.explained
@dirrm lib/rwhois/attribute_defs/
@dirrm lib/rwhois/bin/
@dirrm lib/rwhois/data/asn/
@dirrm lib/rwhois/data/domain/
@dirrm lib/rwhois/data/host/
@dirrm lib/rwhois/data/network/
@dirrm lib/rwhois/data/referral/
@dirrm lib/rwhois/data/soa/
@dirrm lib/rwhois/data/user/
@dirrm lib/rwhois/data/
@dirrm lib/rwhois/etc/
@dirrm lib/rwhois/index_scripts/
@dirrm lib/rwhois/output_templates/
@dirrm lib/rwhois/sbin/
@dirrm lib/rwhois/tmp/
@dirrm lib/rwhois/usr/lib/
@dirrm lib/rwhois/usr/libexec/
@dirrm lib/rwhois/usr/
@dirrm lib/rwhois/
@dirrm share/doc/rwhois
