# dns/openmdns/files/mdnsd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: mdnsd
# REQUIRE: NETWORKING
# BEFORE:  DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name="mdnsd"
desc="%%COMMENT%%"
rcvar="${name}_enable"
command="%%PREFIX%%/sbin/${name}"
extra_commands="reload"

: ${mdnsd_enable:="NO"}
: ${mdnsd_flags:="lo0"}

load_rc_config $name
run_rc_command "$1"
