# net/dhcprelay/files/dhcprelay.in
#! /bin/sh
#
# $FreeBSD$
#
# PROVIDE: dhcprelay
# REQUIRE: DAEMON
#
# Add the following line to /etc/rc.conf to enable dhcrelay:
#
# dhcprelay_enable="YES"
#
# dhcprelay_server	dhcprelay server(s)
# dhcprelay_ifaces	ethernet interface(s)

. /etc/rc.subr

name=dhcprelay
rcvar=${name}_enable

load_rc_config ${name}

: ${dhcprelay_enable="NO"}
: ${dhcprelay_ifaces=""}

pidfile=/var/run/${name}.pid
command=%%PREFIX%%/bin/${name}
command_args="${dhcprelay_ifaces} ${dhcprelay_server}"

run_rc_command "$1"
