# net/dhcprelay/files/openrc-dhcprelay.in
#!/sbin/openrc-run
#
# dhcprelay_server	dhcprelay server(s)
# dhcprelay_ifaces	ethernet interface(s)

name=dhcprelay

: ${dhcprelay_ifaces=""}

pidfile=/var/run/${name}.pid
command=%%PREFIX%%/bin/${name}
command_args="${dhcprelay_ifaces} ${dhcprelay_server}"
