# net/samba-devel/pkg-plist
bin/addtosmbpass
bin/convert_smbpasswd
bin/make_smbpasswd
bin/nmblookup
bin/pdbedit
bin/rpcclient
bin/smbcacls
bin/smbclient
bin/smbcontrol
bin/smbgroupedit
bin/smbpasswd
bin/smbspool
bin/smbstatus
bin/smbtar
bin/smbtree
bin/testparm
bin/testprns
etc/rc.d/samba.sh.sample
etc/smb.conf.default
sbin/nmbd
sbin/smbd
sbin/swat
%%PORTDOCS%%share/doc/samba/COPYING
%%PORTDOCS%%share/doc/samba/Manifest
%%PORTDOCS%%share/doc/samba/README
%%PORTDOCS%%share/doc/samba/README.FreeBSD
%%PORTDOCS%%share/doc/samba/README.swat
%%PORTDOCS%%share/doc/samba/Read-Manifest-Now
%%PORTDOCS%%share/doc/samba/Registry/NT4-Locking.reg
%%PORTDOCS%%share/doc/samba/Registry/NT4_PlainPassword.reg
%%PORTDOCS%%share/doc/samba/Registry/Win2000_PlainPassword.reg
%%PORTDOCS%%share/doc/samba/Registry/Win95_PlainPassword.reg
%%PORTDOCS%%share/doc/samba/Registry/Win98_PlainPassword.reg
%%PORTDOCS%%share/doc/samba/Registry/Win9X-CacheHandling.reg
%%PORTDOCS%%share/doc/samba/Registry/WinME_PlainPassword.reg
%%PORTDOCS%%share/doc/samba/Registry/WinXP_SignOrSeal.reg
%%PORTDOCS%%share/doc/samba/Registry/WindowsTerminalServer.reg
%%PORTDOCS%%share/doc/samba/Roadmap
%%PORTDOCS%%share/doc/samba/THANKS
%%PORTDOCS%%share/doc/samba/WHATSNEW.txt
%%PORTDOCS%%share/doc/samba/announce
%%PORTDOCS%%share/doc/samba/faq/Samba-Server-FAQ-1.html
%%PORTDOCS%%share/doc/samba/faq/Samba-Server-FAQ-2.html
%%PORTDOCS%%share/doc/samba/faq/Samba-Server-FAQ.html
%%PORTDOCS%%share/doc/samba/faq/Samba-Server-FAQ.sgml
%%PORTDOCS%%share/doc/samba/faq/Samba-meta-FAQ-1.html
%%PORTDOCS%%share/doc/samba/faq/Samba-meta-FAQ-2.html
%%PORTDOCS%%share/doc/samba/faq/Samba-meta-FAQ-3.html
%%PORTDOCS%%share/doc/samba/faq/Samba-meta-FAQ-4.html
%%PORTDOCS%%share/doc/samba/faq/Samba-meta-FAQ-5.html
%%PORTDOCS%%share/doc/samba/faq/Samba-meta-FAQ-6.html
%%PORTDOCS%%share/doc/samba/faq/Samba-meta-FAQ.html
%%PORTDOCS%%share/doc/samba/faq/Samba-meta-FAQ.sgml
%%PORTDOCS%%share/doc/samba/faq/Samba-meta-FAQ.txt
%%PORTDOCS%%share/doc/samba/faq/sambafaq-1.html
%%PORTDOCS%%share/doc/samba/faq/sambafaq-2.html
%%PORTDOCS%%share/doc/samba/faq/sambafaq-3.html
%%PORTDOCS%%share/doc/samba/faq/sambafaq-4.html
%%PORTDOCS%%share/doc/samba/faq/sambafaq-5.html
%%PORTDOCS%%share/doc/samba/faq/sambafaq.html
%%PORTDOCS%%share/doc/samba/faq/sambafaq.sgml
%%PORTDOCS%%share/doc/samba/faq/sambafaq.txt
%%PORTDOCS%%share/doc/samba/history
%%PORTDOCS%%share/doc/samba/htmldocs/CVS-Access.html
%%PORTDOCS%%share/doc/samba/htmldocs/DOMAIN_MEMBER.html
%%PORTDOCS%%share/doc/samba/htmldocs/ENCRYPTION.html
%%PORTDOCS%%share/doc/samba/htmldocs/Integrating-with-Windows.html
%%PORTDOCS%%share/doc/samba/htmldocs/NT_Security.html
%%PORTDOCS%%share/doc/samba/htmldocs/OS2-Client-HOWTO.html
%%PORTDOCS%%share/doc/samba/htmldocs/PAM-Authentication-And-Samba.html
%%PORTDOCS%%share/doc/samba/htmldocs/Samba-HOWTO-Collection.html
%%PORTDOCS%%share/doc/samba/htmldocs/Samba-PDC-HOWTO.html
%%PORTDOCS%%share/doc/samba/htmldocs/Samba-LDAP-HOWTO.html
%%PORTDOCS%%share/doc/samba/htmldocs/UNIX_INSTALL.html
%%PORTDOCS%%share/doc/samba/htmldocs/findsmb.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/lmhosts.5.html
%%PORTDOCS%%share/doc/samba/htmldocs/make_smbcodepage.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/make_unicodemap.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/msdfs_setup.html
%%PORTDOCS%%share/doc/samba/htmldocs/net.8.html
%%PORTDOCS%%share/doc/samba/htmldocs/nmbd.8.html
%%PORTDOCS%%share/doc/samba/htmldocs/nmblookup.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/pdbedit.8.html
%%PORTDOCS%%share/doc/samba/htmldocs/printer_driver2.html
%%PORTDOCS%%share/doc/samba/htmldocs/rpcclient.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/samba.7.html
%%PORTDOCS%%share/doc/samba/htmldocs/smb.conf.5.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbcacls.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbclient.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbcontrol.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbd.8.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbgroupedit.8.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbmnt.8.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbmount.8.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbpasswd.5.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbpasswd.8.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbsh.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbspool.8.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbstatus.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbtar.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/smbumount.8.html
%%PORTDOCS%%share/doc/samba/htmldocs/swat.8.html
%%PORTDOCS%%share/doc/samba/htmldocs/testparm.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/testprns.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/wbinfo.1.html
%%PORTDOCS%%share/doc/samba/htmldocs/wfw_slip.htm
%%PORTDOCS%%share/doc/samba/htmldocs/winbind.html
%%PORTDOCS%%share/doc/samba/htmldocs/winbindd.8.html
%%PORTDOCS%%share/doc/samba/textdocs/ADS-HOWTO.txt
%%PORTDOCS%%share/doc/samba/textdocs/Application_Serving.txt
%%PORTDOCS%%share/doc/samba/textdocs/BROWSING-Config.txt
%%PORTDOCS%%share/doc/samba/textdocs/BROWSING.txt
%%PORTDOCS%%share/doc/samba/textdocs/BUGS.txt
%%PORTDOCS%%share/doc/samba/textdocs/DHCP-Server-Configuration.txt
%%PORTDOCS%%share/doc/samba/textdocs/DIAGNOSIS.txt
%%PORTDOCS%%share/doc/samba/textdocs/DNIX.txt
%%PORTDOCS%%share/doc/samba/textdocs/Faxing.txt
%%PORTDOCS%%share/doc/samba/textdocs/GOTCHAS.txt
%%PORTDOCS%%share/doc/samba/textdocs/GROUP-MAPPING-HOWTO.txt
%%PORTDOCS%%share/doc/samba/textdocs/HINTS.txt
%%PORTDOCS%%share/doc/samba/textdocs/INSTALL.sambatar
%%PORTDOCS%%share/doc/samba/textdocs/Imprints.txt
%%PORTDOCS%%share/doc/samba/textdocs/Macintosh_Clients.txt
%%PORTDOCS%%share/doc/samba/textdocs/NetBIOS.txt
%%PORTDOCS%%share/doc/samba/textdocs/PROFILES.txt
%%PORTDOCS%%share/doc/samba/textdocs/Passwords.txt
%%PORTDOCS%%share/doc/samba/textdocs/Printing.txt
%%PORTDOCS%%share/doc/samba/textdocs/README.DCEDFS
%%PORTDOCS%%share/doc/samba/textdocs/README.NOW
%%PORTDOCS%%share/doc/samba/textdocs/README.jis
%%PORTDOCS%%share/doc/samba/textdocs/README.sambatar
%%PORTDOCS%%share/doc/samba/textdocs/Recent-FAQs.txt
%%PORTDOCS%%share/doc/samba/textdocs/RoutedNetworks.txt
%%PORTDOCS%%share/doc/samba/textdocs/SCO.txt
%%PORTDOCS%%share/doc/samba/textdocs/SMBTAR.notes
%%PORTDOCS%%share/doc/samba/textdocs/Samba-OpenSSL.txt
%%PORTDOCS%%share/doc/samba/textdocs/Speed.txt
%%PORTDOCS%%share/doc/samba/textdocs/Speed2.txt
%%PORTDOCS%%share/doc/samba/textdocs/Tracing.txt
%%PORTDOCS%%share/doc/samba/textdocs/UNIX-SMB.txt
%%PORTDOCS%%share/doc/samba/textdocs/UNIX_SECURITY.txt
%%PORTDOCS%%share/doc/samba/textdocs/Win95.txt
%%PORTDOCS%%share/doc/samba/textdocs/WinNT.txt
%%PORTDOCS%%share/doc/samba/textdocs/cifsntdomain.txt
%%PORTDOCS%%share/doc/samba/textdocs/kurs.pdf
%%PORTDOCS%%share/doc/samba/textdocs/kurs.tex
%%PORTDOCS%%share/doc/samba/textdocs/logo.ps
%%PORTDOCS%%share/doc/samba/textdocs/security_level.txt
share/examples/samba/README
share/examples/samba/VFS/.cvsignore
share/examples/samba/VFS/Makefile
share/examples/samba/VFS/README
share/examples/samba/VFS/block/Makefile
share/examples/samba/VFS/block/block.c
share/examples/samba/VFS/block/samba-block.conf
share/examples/samba/VFS/block/smb.conf
share/examples/samba/VFS/audit.c
share/examples/samba/VFS/skel.c
share/examples/samba/LDAP/README
share/examples/samba/LDAP/export2_smbpasswd.pl
share/examples/samba/LDAP/export_smbpasswd.pl
share/examples/samba/LDAP/import2_smbpasswd.pl
share/examples/samba/LDAP/import_smbpasswd.pl
share/examples/samba/LDAP/ldapchpasswd
share/examples/samba/LDAP/ldapsync.pl
share/examples/samba/LDAP/samba.schema
share/examples/samba/appliance/Makefile
share/examples/samba/appliance/README
share/examples/samba/appliance/appliance.spec
share/examples/samba/appliance/build.sh
share/examples/samba/appliance/smb.conf-appliance
share/examples/samba/autofs/auto.a
share/examples/samba/dce-dfs/README
share/examples/samba/dce-dfs/smb.conf
share/examples/samba/libsmbclient/Makefile
share/examples/samba/libsmbclient/README
share/examples/samba/libsmbclient/testsmbc.c
share/examples/samba/libsmbclient/tree.c
share/examples/samba/misc/extra_smbstatus
share/examples/samba/misc/swat.pl
share/examples/samba/misc/wall.perl
share/examples/samba/printer-accounting/README
share/examples/samba/printer-accounting/acct-all
share/examples/samba/printer-accounting/acct-sum
share/examples/samba/printer-accounting/hp5-redir
share/examples/samba/printer-accounting/lp-acct
share/examples/samba/printer-accounting/printcap
share/examples/samba/printing/smbprint
share/examples/samba/printing/smbprint.sysv
share/examples/samba/simple/README
share/examples/samba/simple/smb.conf
share/examples/samba/smb.conf.default
share/examples/samba/svr4-startup/README
share/examples/samba/svr4-startup/samba.server
share/examples/samba/thoralf/smb.conf
share/examples/samba/tridge/README
share/examples/samba/tridge/smb.conf
share/examples/samba/tridge/smb.conf.WinNT
share/examples/samba/tridge/smb.conf.fjall
share/examples/samba/tridge/smb.conf.lapland
share/examples/samba/tridge/smb.conf.vittjokk
share/examples/samba/validchars/msdos70.out
share/examples/samba/validchars/nwdos70.out
share/examples/samba/validchars/readme
share/examples/samba/validchars/validchr.c
share/examples/samba/validchars/validchr.com
share/swat/help/CVS-Access.html
share/swat/help/DOMAIN_MEMBER.html
share/swat/help/ENCRYPTION.html
share/swat/help/Integrating-with-Windows.html
share/swat/help/NT_Security.html
share/swat/help/OS2-Client-HOWTO.html
share/swat/help/PAM-Authentication-And-Samba.html
share/swat/help/Samba-HOWTO-Collection.html
share/swat/help/Samba-PDC-HOWTO.html
share/swat/help/UNIX_INSTALL.html
share/swat/help/findsmb.1.html
share/swat/help/lmhosts.5.html
share/swat/help/make_smbcodepage.1.html
share/swat/help/make_unicodemap.1.html
share/swat/help/msdfs_setup.html
share/swat/help/net.8.html
share/swat/help/nmbd.8.html
share/swat/help/nmblookup.1.html
share/swat/help/pdbedit.8.html
share/swat/help/printer_driver2.html
share/swat/help/rpcclient.1.html
share/swat/help/Samba-LDAP-HOWTO.html
share/swat/help/samba.7.html
share/swat/help/smb.conf.5.html
share/swat/help/smbcacls.1.html
share/swat/help/smbclient.1.html
share/swat/help/smbcontrol.1.html
share/swat/help/smbd.8.html
share/swat/help/smbgroupedit.8.html
share/swat/help/smbmnt.8.html
share/swat/help/smbmount.8.html
share/swat/help/smbpasswd.5.html
share/swat/help/smbpasswd.8.html
share/swat/help/smbsh.1.html
share/swat/help/smbspool.8.html
share/swat/help/smbstatus.1.html
share/swat/help/smbtar.1.html
share/swat/help/smbumount.8.html
share/swat/help/swat.8.html
share/swat/help/testparm.1.html
share/swat/help/testprns.1.html
share/swat/help/wbinfo.1.html
share/swat/help/welcome.html
share/swat/help/winbind.html
share/swat/help/winbindd.8.html
share/swat/images/globals.gif
share/swat/images/home.gif
share/swat/images/passwd.gif
share/swat/images/printers.gif
share/swat/images/samba.gif
share/swat/images/shares.gif
share/swat/images/status.gif
share/swat/images/viewconfig.gif
share/swat/include/footer.html
share/swat/include/header.html
share/swat/lang/ja/help/welcome.html
share/swat/lang/ja/images/globals.gif
share/swat/lang/ja/images/home.gif
share/swat/lang/ja/images/passwd.gif
share/swat/lang/ja/images/printers.gif
share/swat/lang/ja/images/samba.gif
share/swat/lang/ja/images/shares.gif
share/swat/lang/ja/images/status.gif
share/swat/lang/ja/images/viewconfig.gif
share/swat/lang/ja/include/footer.html
share/swat/lang/ja/include/header.html
share/swat/lang/ja/include/header.nocss.html
share/swat/lang/ja/include/header_css.html
share/swat/lang/tr/help/welcome.html
share/swat/lang/tr/images/globals.gif
share/swat/lang/tr/images/home.gif
share/swat/lang/tr/images/passwd.gif
share/swat/lang/tr/images/printers.gif
share/swat/lang/tr/images/samba.gif
share/swat/lang/tr/images/shares.gif
share/swat/lang/tr/images/status.gif
share/swat/lang/tr/images/viewconfig.gif
share/swat/lang/tr/include/header.html
share/swat/using_samba/appa_01.html
share/swat/using_samba/appa_02.html
share/swat/using_samba/appa_03.html
share/swat/using_samba/appa_04.html
share/swat/using_samba/appa_05.html
share/swat/using_samba/appb_01.html
share/swat/using_samba/appb_02.html
share/swat/using_samba/appb_03.html
share/swat/using_samba/appc_01.html
share/swat/using_samba/appd_01.html
share/swat/using_samba/appe_01.html
share/swat/using_samba/appf_01.html
share/swat/using_samba/ch01_01.html
share/swat/using_samba/ch01_02.html
share/swat/using_samba/ch01_03.html
share/swat/using_samba/ch01_04.html
share/swat/using_samba/ch01_05.html
share/swat/using_samba/ch01_06.html
share/swat/using_samba/ch01_07.html
share/swat/using_samba/ch01_08.html
share/swat/using_samba/ch02_01.html
share/swat/using_samba/ch02_02.html
share/swat/using_samba/ch02_03.html
share/swat/using_samba/ch02_04.html
share/swat/using_samba/ch02_05.html
share/swat/using_samba/ch02_06.html
share/swat/using_samba/ch03_01.html
share/swat/using_samba/ch03_02.html
share/swat/using_samba/ch03_03.html
share/swat/using_samba/ch04_01.html
share/swat/using_samba/ch04_02.html
share/swat/using_samba/ch04_03.html
share/swat/using_samba/ch04_04.html
share/swat/using_samba/ch04_05.html
share/swat/using_samba/ch04_06.html
share/swat/using_samba/ch04_07.html
share/swat/using_samba/ch04_08.html
share/swat/using_samba/ch05_01.html
share/swat/using_samba/ch05_02.html
share/swat/using_samba/ch05_03.html
share/swat/using_samba/ch05_04.html
share/swat/using_samba/ch05_05.html
share/swat/using_samba/ch06_01.html
share/swat/using_samba/ch06_02.html
share/swat/using_samba/ch06_03.html
share/swat/using_samba/ch06_04.html
share/swat/using_samba/ch06_05.html
share/swat/using_samba/ch06_06.html
share/swat/using_samba/ch07_01.html
share/swat/using_samba/ch07_02.html
share/swat/using_samba/ch07_03.html
share/swat/using_samba/ch08_01.html
share/swat/using_samba/ch08_02.html
share/swat/using_samba/ch08_03.html
share/swat/using_samba/ch08_04.html
share/swat/using_samba/ch08_05.html
share/swat/using_samba/ch08_06.html
share/swat/using_samba/ch08_07.html
share/swat/using_samba/ch09_01.html
share/swat/using_samba/ch09_02.html
share/swat/using_samba/ch09_03.html
share/swat/using_samba/figs/sam.0101.gif
share/swat/using_samba/figs/sam.0102.gif
share/swat/using_samba/figs/sam.0103.gif
share/swat/using_samba/figs/sam.0104.gif
share/swat/using_samba/figs/sam.0105.gif
share/swat/using_samba/figs/sam.0106.gif
share/swat/using_samba/figs/sam.0107.gif
share/swat/using_samba/figs/sam.0108.gif
share/swat/using_samba/figs/sam.0109.gif
share/swat/using_samba/figs/sam.0110.gif
share/swat/using_samba/figs/sam.0111.gif
share/swat/using_samba/figs/sam.0112.gif
share/swat/using_samba/figs/sam.0113.gif
share/swat/using_samba/figs/sam.0114.gif
share/swat/using_samba/figs/sam.0201.gif
share/swat/using_samba/figs/sam.0202.gif
share/swat/using_samba/figs/sam.0203.gif
share/swat/using_samba/figs/sam.0204.gif
share/swat/using_samba/figs/sam.0301.gif
share/swat/using_samba/figs/sam.0302.gif
share/swat/using_samba/figs/sam.0303.gif
share/swat/using_samba/figs/sam.0304.gif
share/swat/using_samba/figs/sam.0305.gif
share/swat/using_samba/figs/sam.0306.gif
share/swat/using_samba/figs/sam.0307.gif
share/swat/using_samba/figs/sam.0308.gif
share/swat/using_samba/figs/sam.0309.gif
share/swat/using_samba/figs/sam.0310.gif
share/swat/using_samba/figs/sam.0311.gif
share/swat/using_samba/figs/sam.0312.gif
share/swat/using_samba/figs/sam.0313.gif
share/swat/using_samba/figs/sam.0314.gif
share/swat/using_samba/figs/sam.0315.gif
share/swat/using_samba/figs/sam.0316.gif
share/swat/using_samba/figs/sam.0317.gif
share/swat/using_samba/figs/sam.0318.gif
share/swat/using_samba/figs/sam.0319.gif
share/swat/using_samba/figs/sam.0320.gif
share/swat/using_samba/figs/sam.0321.gif
share/swat/using_samba/figs/sam.0322.gif
share/swat/using_samba/figs/sam.0323.gif
share/swat/using_samba/figs/sam.0324.gif
share/swat/using_samba/figs/sam.0325.gif
share/swat/using_samba/figs/sam.0326.gif
share/swat/using_samba/figs/sam.0327.gif
share/swat/using_samba/figs/sam.0328.gif
share/swat/using_samba/figs/sam.0401.gif
share/swat/using_samba/figs/sam.0402.gif
share/swat/using_samba/figs/sam.0403.gif
share/swat/using_samba/figs/sam.0404.gif
share/swat/using_samba/figs/sam.0405.gif
share/swat/using_samba/figs/sam.0406.gif
share/swat/using_samba/figs/sam.0407.gif
share/swat/using_samba/figs/sam.0501.gif
share/swat/using_samba/figs/sam.0502.gif
share/swat/using_samba/figs/sam.0503.gif
share/swat/using_samba/figs/sam.0504.gif
share/swat/using_samba/figs/sam.0505.gif
share/swat/using_samba/figs/sam.0506.gif
share/swat/using_samba/figs/sam.0507.gif
share/swat/using_samba/figs/sam.0508.gif
share/swat/using_samba/figs/sam.0601.gif
share/swat/using_samba/figs/sam.0602.gif
share/swat/using_samba/figs/sam.0603.gif
share/swat/using_samba/figs/sam.0604.gif
share/swat/using_samba/figs/sam.0605.gif
share/swat/using_samba/figs/sam.0606.gif
share/swat/using_samba/figs/sam.0701.gif
share/swat/using_samba/figs/sam.0702.gif
share/swat/using_samba/figs/sam.0703.gif
share/swat/using_samba/figs/sam.0704.gif
share/swat/using_samba/figs/sam.0705.gif
share/swat/using_samba/figs/sam.0706.gif
share/swat/using_samba/figs/sam.0707.gif
share/swat/using_samba/figs/sam.0708.gif
share/swat/using_samba/figs/sam.0709.gif
share/swat/using_samba/figs/sam.0801.gif
share/swat/using_samba/figs/sam.0802.gif
share/swat/using_samba/figs/sam.0803.gif
share/swat/using_samba/figs/sam.0804.gif
share/swat/using_samba/figs/sam.0805.gif
share/swat/using_samba/figs/sam.0901.gif
share/swat/using_samba/figs/sam.0902.gif
share/swat/using_samba/figs/sam.0903.gif
share/swat/using_samba/figs/sam.0904.gif
share/swat/using_samba/figs/sam.0905.gif
share/swat/using_samba/figs/sam.aa01.gif
share/swat/using_samba/figs/sam.ab01.gif
share/swat/using_samba/figs/sam.ab02.gif
share/swat/using_samba/gifs/index.gif
share/swat/using_samba/gifs/samba.s.gif
share/swat/using_samba/gifs/txthome.gif
share/swat/using_samba/gifs/txtnexta.gif
share/swat/using_samba/gifs/txtpreva.gif
share/swat/using_samba/index.html
share/swat/using_samba/inx.html
share/swat/using_samba/licenseinfo.html
share/swat/using_samba/this_edition.html
@dirrm share/swat/using_samba/gifs
@dirrm share/swat/using_samba/figs
@dirrm share/swat/using_samba
@dirrm share/swat/include
@dirrm share/swat/images
@dirrm share/swat/help
@dirrm share/swat/lang/ja/include
@dirrm share/swat/lang/ja/images
@dirrm share/swat/lang/ja/help
@dirrm share/swat/lang/ja
@dirrm share/swat/lang/tr/include
@dirrm share/swat/lang/tr/images
@dirrm share/swat/lang/tr/help
@dirrm share/swat/lang/tr
@dirrm share/swat/lang
@dirrm share/swat
@dirrm share/examples/samba/LDAP
@dirrm share/examples/samba/validchars
@dirrm share/examples/samba/tridge
@dirrm share/examples/samba/thoralf
@dirrm share/examples/samba/svr4-startup
@dirrm share/examples/samba/simple
@dirrm share/examples/samba/printing
@dirrm share/examples/samba/printer-accounting
@dirrm share/examples/samba/misc
@dirrm share/examples/samba/libsmbclient
@dirrm share/examples/samba/dce-dfs
@dirrm share/examples/samba/autofs
@dirrm share/examples/samba/appliance
@dirrm share/examples/samba/VFS/block
@dirrm share/examples/samba/VFS
@dirrm share/examples/samba
%%PORTDOCS%%@dirrm share/doc/samba/Registry
%%PORTDOCS%%@dirrm share/doc/samba/textdocs
%%PORTDOCS%%@dirrm share/doc/samba/htmldocs
%%PORTDOCS%%@dirrm share/doc/samba/faq
%%PORTDOCS%%@dirrm share/doc/samba
