# chinese/FreeWnn-server/pkg-plist
@cwd /usr/local
etc/rc.d/cWnn.sh
lib/wnn/zh_CN/cserverrc
lib/wnn/zh_CN/cixing.data
lib/wnn/zh_CN/uum.msg
lib/wnn/zh_CN/uumkey
lib/wnn/zh_CN/uumkey_e
lib/wnn/zh_CN/uumkey_p
lib/wnn/zh_CN/uumrc
lib/wnn/zh_CN/uumrc_p
lib/wnn/zh_CN/uumrc_z
lib/wnn/zh_CN/wnnstat.msg
lib/wnn/zh_CN/libwnn.msg
lib/wnn/zh_CN/wnnenvrc
lib/wnn/zh_CN/wnnenvrc_Qi
lib/wnn/zh_CN/wnnenvrc_QiR
lib/wnn/zh_CN/wnnenvrc_R
lib/wnn/zh_CN/wnnenvrc_Wu
lib/wnn/zh_CN/wnnenvrc_WuR
lib/wnn/zh_CN/rk/1B_BS
lib/wnn/zh_CN/rk/1B_TOLOWER
lib/wnn/zh_CN/rk/1B_TOUPPER
lib/wnn/zh_CN/rk/1Z_ZhuYin
lib/wnn/zh_CN/rk/2A_CTRL
lib/wnn/zh_CN/rk/2B_WuBi
lib/wnn/zh_CN/rk/2P_QuanPin
lib/wnn/zh_CN/rk/2P_RongCuo
lib/wnn/zh_CN/rk/2Z_ZhuYin
lib/wnn/zh_CN/rk/2Z_tail
lib/wnn/zh_CN/rk/2Z_tail_ma
lib/wnn/zh_CN/rk/2Z_tail_pin
lib/wnn/zh_CN/rk/3B_quanjiao
lib/wnn/zh_CN/rk/mode
lib/wnn/zh_CN/rk_p/1B_BS
lib/wnn/zh_CN/rk_p/1B_TOLOWER
lib/wnn/zh_CN/rk_p/1B_TOUPPER
lib/wnn/zh_CN/rk_p/1F_select
lib/wnn/zh_CN/rk_p/2A_CTRL
lib/wnn/zh_CN/rk_p/2B_QianMa
lib/wnn/zh_CN/rk_p/2B_WuBi
lib/wnn/zh_CN/rk_p/2B_ch_yo_zi
lib/wnn/zh_CN/rk_p/2B_ma_GuoBi
lib/wnn/zh_CN/rk_p/2B_ma_NeiMa
lib/wnn/zh_CN/rk_p/2B_ma_QuWei
lib/wnn/zh_CN/rk_p/2F_restart
lib/wnn/zh_CN/rk_p/2F_select
lib/wnn/zh_CN/rk_p/2P_ErPin
lib/wnn/zh_CN/rk_p/2P_ErPin_3K
lib/wnn/zh_CN/rk_p/2P_ErPin_P
lib/wnn/zh_CN/rk_p/2P_ErPin_P3K
lib/wnn/zh_CN/rk_p/2P_QuanPin
lib/wnn/zh_CN/rk_p/2P_RC2
lib/wnn/zh_CN/rk_p/2P_RC2_P
lib/wnn/zh_CN/rk_p/2P_RongCuo
lib/wnn/zh_CN/rk_p/2Z_tail_ma
lib/wnn/zh_CN/rk_p/2Z_tail_pin
lib/wnn/zh_CN/rk_p/3B_quanjiao
lib/wnn/zh_CN/rk_p/3P_ErPin
lib/wnn/zh_CN/rk_p/3P_ErPin_P
lib/wnn/zh_CN/rk_p/mode
lib/wnn/zh_CN/rk_z/1B_TOUPPER
lib/wnn/zh_CN/rk_z/1F_select
lib/wnn/zh_CN/rk_z/1Z_ZhuYin
lib/wnn/zh_CN/rk_z/1Z_ZhuYin_X
lib/wnn/zh_CN/rk_z/2A_CTRL
lib/wnn/zh_CN/rk_z/2F_select
lib/wnn/zh_CN/rk_z/2Z_ZhuYin
lib/wnn/zh_CN/rk_z/2Z_ZhuYin_X
lib/wnn/zh_CN/rk_z/2Z_tail
lib/wnn/zh_CN/rk_z/2Z_tail_X
lib/wnn/zh_CN/rk_z/3B_quanjiao
lib/wnn/zh_CN/rk_z/mode
lib/wnn/zh_CN/dic/sys/full.con
lib/wnn/zh_CN/dic/sys/full.conR
lib/wnn/zh_CN/dic/sys/basic.dic
lib/wnn/zh_CN/dic/sys/level_1.dic
lib/wnn/zh_CN/dic/sys/level_2.dic
lib/wnn/zh_CN/dic/sys/computer.dic
lib/wnn/zh_CN/dic/sys/cwnn.dic
lib/wnn/zh_CN/dic/sys/WuBi.dic
lib/wnn/zh_CN/dic/sys/QianMa.dic
lib/wnn/cvt_key_empty
lib/wnn/cvt_key_tbl
lib/wnn/cvt_key_tbl.ST
lib/wnn/cvt_key_tbl.gm
lib/wnn/cvt_key_tbl.kt
lib/wnn/cvt_key_tbl.mv
lib/wnn/cvt_key_tbl.vt
lib/wnn/serverdefs
lib/wnn/zh_TW/dic/sys/full.con
lib/wnn/zh_TW/dic/sys/full.conR
lib/wnn/zh_TW/dic/sys/cns_ch.dic
lib/wnn/zh_TW/dic/sys/cns_wd.dic
lib/wnn/zh_TW/tserverrc
lib/wnn/zh_TW/cixing.data
lib/wnn/zh_TW/uum.msg
lib/wnn/zh_TW/uumkey
lib/wnn/zh_TW/uumkey_e
lib/wnn/zh_TW/uumkey_p
lib/wnn/zh_TW/uumrc
lib/wnn/zh_TW/uumrc_p
lib/wnn/zh_TW/uumrc_z
lib/wnn/zh_TW/wnnstat.msg
lib/wnn/zh_TW/libwnn.msg
lib/wnn/zh_TW/wnnenvrc
lib/wnn/zh_TW/wnnenvrc_R
lib/wnn/zh_TW/rk/mode
lib/wnn/zh_TW/rk/2A_CTRL
lib/wnn/zh_TW/rk_p/mode
lib/wnn/zh_TW/rk_p/2A_CTRL
lib/wnn/zh_TW/rk_z/mode
lib/wnn/zh_TW/rk_z/2A_CTRL
@comment Xwnmo
lib/wnn/ja_JP/uumrc_vi.xim
lib/wnn/ja_JP/uumrc_e.xim
lib/wnn/ja_JP/uumrc.xim
lib/wnn/ja_JP/uumkey.xim
lib/wnn/ja_JP/uumkey_v.xim
lib/wnn/ja_JP/uumkey_e.xim
lib/wnn/ja_JP/xim.msg
lib/wnn/ja_JP/xjutil.msg
lib/wnn/ko_KR/uumkey.xim
lib/wnn/ko_KR/uumrc.xim
lib/wnn/ko_KR/xim.msg
lib/wnn/ko_KR/xjutil.msg
lib/wnn/lt_LN/uumkey.xim
lib/wnn/lt_LN/uumrc.xim
lib/wnn/lt_LN/xim.msg
lib/wnn/zh_CN/uumkey.xim
lib/wnn/zh_CN/uumkey_p.xim
lib/wnn/zh_CN/uumrc.xim
lib/wnn/zh_CN/uumrc_p.xim
lib/wnn/zh_CN/uumrc_z.xim
lib/wnn/zh_CN/xim.msg
lib/wnn/zh_CN/xjutil.msg
lib/wnn/zh_TW/uumkey.xim
lib/wnn/zh_TW/uumkey_p.xim
lib/wnn/zh_TW/uumrc.xim
lib/wnn/zh_TW/uumrc_p.xim
lib/wnn/zh_TW/uumrc_z.xim
lib/wnn/zh_TW/xim.msg
lib/wnn/zh_TW/xjutil.msg
lib/wnn/cvt_xim_tbl
lib/wnn/ximconf
lib/wnn/ximrc
lib/wnn/ximrc_vi
@comment END Xwnmo
@dirrm lib/wnn/zh_CN/rk_z
@dirrm lib/wnn/zh_CN/rk_p
@dirrm lib/wnn/zh_CN/rk
@dirrm lib/wnn/zh_CN/dic/sys
@dirrm lib/wnn/zh_CN/dic/usr
@dirrm lib/wnn/zh_CN/dic
@dirrm lib/wnn/zh_CN
@dirrm lib/wnn/zh_TW/rk_z
@dirrm lib/wnn/zh_TW/rk_p
@dirrm lib/wnn/zh_TW/rk
@dirrm lib/wnn/zh_TW/dic/sys
@dirrm lib/wnn/zh_TW/dic/usr
@dirrm lib/wnn/zh_TW/dic
@dirrm lib/wnn/zh_TW
@comment Xwnmo
@dirrm lib/wnn/ja_JP
@dirrm lib/wnn/lt_LN
@dirrm lib/wnn/ko_KR
@comment END Xwnmo
@dirrm lib/wnn
bin/cWnn4/cserver
bin/cWnn4/tserver
bin/cWnn4/cuum
bin/cWnn4/cwnnstat
bin/cWnn4/cwnnkill
bin/cWnn4/cwddel
bin/cWnn4/cwdreg
bin/cWnn4/cdtoa
bin/cWnn4/catod
bin/cWnn4/catof
bin/cWnn4/cwnntouch
@dirrm bin/cWnn4
share/doc/zh-Wnn/manual/README
share/doc/zh-Wnn/manual/app
share/doc/zh-Wnn/manual/chap1
share/doc/zh-Wnn/manual/chap2
share/doc/zh-Wnn/manual/chap3
share/doc/zh-Wnn/manual/chap4
share/doc/zh-Wnn/manual/chap5
share/doc/zh-Wnn/manual/chap6
share/doc/zh-Wnn/manual/chap7
share/doc/zh-Wnn/manual/chap8
share/doc/zh-Wnn/manual/contents
share/doc/zh-Wnn/manual/glossary
share/doc/zh-Wnn/manual/index
share/doc/zh-Wnn/manual/preface
share/doc/zh-Wnn/manual.en/README
share/doc/zh-Wnn/manual.en/Table-a-2.1.Z
share/doc/zh-Wnn/manual.en/Table-a-2.2.Z
share/doc/zh-Wnn/manual.en/Table-a-2.3.Z
share/doc/zh-Wnn/manual.en/Table-a-3.1.Z
share/doc/zh-Wnn/manual.en/Table-a-3.2.Z
share/doc/zh-Wnn/manual.en/Table-a-3.3.Z
share/doc/zh-Wnn/manual.en/Table-a-4.1.Z
share/doc/zh-Wnn/manual.en/Table-a-4.2.Z
share/doc/zh-Wnn/manual.en/Table-a-4.3.Z
share/doc/zh-Wnn/manual.en/Table-a-4.4.Z
share/doc/zh-Wnn/manual.en/Table-a-6.1.Z
share/doc/zh-Wnn/manual.en/Table-a-6.2.Z
share/doc/zh-Wnn/manual.en/Table-c-2.1.Z
share/doc/zh-Wnn/manual.en/Table-c-2.2.Z
share/doc/zh-Wnn/manual.en/Table-c-2.3.Z
share/doc/zh-Wnn/manual.en/Table-c-2.4.Z
share/doc/zh-Wnn/manual.en/Table-c-2.5.Z
share/doc/zh-Wnn/manual.en/Table-c-2.6.Z
share/doc/zh-Wnn/manual.en/Table-c-2.7.Z
share/doc/zh-Wnn/manual.en/app
share/doc/zh-Wnn/manual.en/chap1
share/doc/zh-Wnn/manual.en/Table-c-2.8.Z
share/doc/zh-Wnn/manual.en/Table-c-3.1.Z
share/doc/zh-Wnn/manual.en/Table-c-3.2.Z
share/doc/zh-Wnn/manual.en/Table-c-3.3.Z
share/doc/zh-Wnn/manual.en/Table-c-3.4.Z
share/doc/zh-Wnn/manual.en/Table-c-3.5.Z
share/doc/zh-Wnn/manual.en/Table-c-5.1.Z
share/doc/zh-Wnn/manual.en/Table-c-6.1.Z
share/doc/zh-Wnn/manual.en/Table-c-6.2.Z
share/doc/zh-Wnn/manual.en/Table-c-6.3.Z
share/doc/zh-Wnn/manual.en/Table-c-6.4.Z
share/doc/zh-Wnn/manual.en/Table-c-6.5.Z
share/doc/zh-Wnn/manual.en/Table-c-6.6.Z
share/doc/zh-Wnn/manual.en/Table-c-6.7.Z
share/doc/zh-Wnn/manual.en/Table-c-6.8.Z
share/doc/zh-Wnn/manual.en/Table-c-7.1.Z
share/doc/zh-Wnn/manual.en/Table-c-7.2.Z
share/doc/zh-Wnn/manual.en/chap2
share/doc/zh-Wnn/manual.en/chap3
share/doc/zh-Wnn/manual.en/chap4
share/doc/zh-Wnn/manual.en/chap5
share/doc/zh-Wnn/manual.en/chap6
share/doc/zh-Wnn/manual.en/chap7
share/doc/zh-Wnn/manual.en/content
share/doc/zh-Wnn/manual.en/cover
share/doc/zh-Wnn/manual.en/index
share/doc/zh-Wnn/manual.en/preface
@dirrm share/doc/zh-Wnn/manual.en
@dirrm share/doc/zh-Wnn/manual
@dirrm share/doc/zh-Wnn
@cwd /usr/X11R6
@comment Xwnmo
bin/xwnmo
bin/killxwnmo
bin/xjutil
lib/X11/app-defaults/Xwnmo
@comment END Xwnmo
lib/libcwnn4.so
lib/libcwnn4.so.1.0
@exec /usr/bin/env OBJFORMAT=%%PORTOBJFORMAT%% /sbin/ldconfig -m %D/lib
@unexec /usr/bin/env OBJFORMAT=%%PORTOBJFORMAT%% /sbin/ldconfig -R
lib/libcwnn.a
include/cwnn/jllib.h
include/cwnn/jslib.h
include/cwnn/rk_spclval.h
include/cwnn/wnnerror.h
include/cwnn/commonhd.h
include/cwnn/cplib.h
include/cwnn/msg.h
@dirrm include/cwnn
man/man1/cwdreg.1.gz
man/man1/cdicsort.1.gz
man/man1/cwddel.1.gz
man/man1/cdtoa.1.gz
man/man1/catof.1.gz
man/man1/catod.1.gz
man/man1/cwnnstat.1.gz
man/man1/cwnnkill.1.gz
man/man1/cuum.1.gz
man/man1/cserver.1.gz
man/man4/cst_end.4.gz
man/man4/cwnn.4.gz
man/man4/cenv.4.gz
man/man4/ckey.4.gz
@comment Xwnmo
man/man1/xwnmo.1.gz
man/man1/ximrc.1.gz
man/man1/ximconf.1.gz
man/man1/cvt_fun_tbl.1.gz
man/man1/cvt_meta_tbl.1.gz
man/man1/killxwnmo.1.gz
man/man1/cvt_xim_tbl.1.gz
man/man1/xjutil.1.gz
@comment END Xwnmo
