# net/dhcp6/pkg-plist
etc/dhcp6s.conf.sample
etc/dhcp6c.conf.sample
etc/rc.d/dhcp6c.sh
etc/rc.d/dhcp6relay.sh
etc/rc.d/dhcp6s.sh
sbin/dhcp6c
sbin/dhcp6relay
sbin/dhcp6s
sbin/dhcp6ctl
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%@dirrm %%DOCSDIR%%
