# www/cocoon/pkg-plist
sbin/cocoon.sh
sbin/cocoonctl
etc/rc.d/cocoon.sh
@dirrm cocoon/tools
@dirrm cocoon
