# www/cocoon/pkg-plist
@comment $FreeBSD$
sbin/cocoonctl
etc/rc.d/cocoon.sh
@dirrm cocoon/tools
@dirrm cocoon
