# www/apache13-modssl/pkg-plist
bin/dbmmanage
bin/htdigest
bin/htpasswd
sbin/ab
sbin/apache
sbin/apachectl
sbin/apxs
sbin/logresolve
sbin/rotatelogs
etc/rc.d/apache.sh
etc/apache/apache.conf
etc/apache/apache.conf.default
etc/apache/magic
etc/apache/magic.default
etc/apache/mime.types
etc/apache/mime.types.default
etc/apache/ssl.crt/README.CRT
etc/apache/ssl.crt/Makefile
etc/apache/ssl.crt/ca-bundle.crt
etc/apache/ssl.crt/snakeoil-ca-rsa.crt
etc/apache/ssl.crt/snakeoil-ca-dsa.crt
etc/apache/ssl.crt/snakeoil-rsa.crt
etc/apache/ssl.crt/snakeoil-dsa.crt
etc/apache/ssl.crt/server.crt
etc/apache/ssl.crl/README.CRL
etc/apache/ssl.crl/Makefile
etc/apache/ssl.csr/README.CSR
etc/apache/ssl.csr/server.csr
etc/apache/ssl.prm/README.PRM
etc/apache/ssl.prm/snakeoil-ca-dsa.prm
etc/apache/ssl.prm/snakeoil-dsa.prm
etc/apache/ssl.key/README.KEY
etc/apache/ssl.key/snakeoil-ca-rsa.key
etc/apache/ssl.key/snakeoil-ca-dsa.key
etc/apache/ssl.key/snakeoil-rsa.key
etc/apache/ssl.key/snakeoil-dsa.key
etc/apache/ssl.key/server.key
include/apache/alloc.h
include/apache/ap_ctx.h
include/apache/ap_hook.h
include/apache/ap_mm.h
include/apache/ap.h
include/apache/ap_compat.h
include/apache/ap_config.h
include/apache/ap_config_auto.h
include/apache/ap_ctype.h
include/apache/ap_md5.h
include/apache/ap_mmn.h
include/apache/ap_sha1.h
include/apache/buff.h
include/apache/compat.h
include/apache/conf.h
include/apache/explain.h
include/apache/fnmatch.h
include/apache/hsregex.h
include/apache/http_conf_globals.h
include/apache/http_config.h
include/apache/http_core.h
include/apache/http_log.h
include/apache/http_main.h
include/apache/http_protocol.h
include/apache/http_request.h
include/apache/http_vhost.h
include/apache/httpd.h
include/apache/multithread.h
include/apache/os.h
include/apache/os-inline.c
include/apache/rfc1413.h
include/apache/scoreboard.h
include/apache/util_date.h
include/apache/util_md5.h
include/apache/util_script.h
include/apache/util_uri.h
include/apache/xml/asciitab.h
include/apache/xml/xmltok.h
include/apache/xml/xmlrole.h
include/apache/xml/xmlparse.h
include/apache/xml/xmldef.h
include/apache/xml/utf8tab.h
include/apache/xml/nametab.h
include/apache/xml/latin1tab.h
include/apache/xml/iasciitab.h
include/apache/xml/hashtable.h
include/apache/xml/xmltok_impl.h
libexec/apache/httpd.exp
libexec/apache/libproxy.so
libexec/apache/libssl.so
libexec/apache/mod_define.so
libexec/apache/mod_access.so
libexec/apache/mod_actions.so
libexec/apache/mod_alias.so
libexec/apache/mod_asis.so
libexec/apache/mod_auth.so
libexec/apache/mod_auth_anon.so
libexec/apache/mod_auth_db.so
libexec/apache/mod_autoindex.so
libexec/apache/mod_cern_meta.so
libexec/apache/mod_cgi.so
libexec/apache/mod_digest.so
libexec/apache/mod_dir.so
libexec/apache/mod_env.so
libexec/apache/mod_expires.so
libexec/apache/mod_headers.so
libexec/apache/mod_imap.so
libexec/apache/mod_include.so
libexec/apache/mod_info.so
libexec/apache/mod_log_config.so
libexec/apache/mod_mime.so
libexec/apache/mod_mime_magic.so
libexec/apache/mod_negotiation.so
libexec/apache/mod_rewrite.so
libexec/apache/mod_setenvif.so
libexec/apache/mod_speling.so
libexec/apache/mod_status.so
libexec/apache/mod_unique_id.so
libexec/apache/mod_userdir.so
libexec/apache/mod_usertrack.so
libexec/apache/mod_vhost_alias.so
share/apache/cgi-bin/printenv
share/apache/cgi-bin/test-cgi
share/apache/htdocs/apache_pb.gif
share/apache/htdocs/index.html
share/apache/htdocs/manual/LICENSE
share/apache/htdocs/manual/bind.html
share/apache/htdocs/manual/cgi_path.html
share/apache/htdocs/manual/content-negotiation.html
share/apache/htdocs/manual/custom-error.html
share/apache/htdocs/manual/dns-caveats.html
share/apache/htdocs/manual/dso.html
share/apache/htdocs/manual/ebcdic.html
share/apache/htdocs/manual/env.html
share/apache/htdocs/manual/footer.html
share/apache/htdocs/manual/handler.html
share/apache/htdocs/manual/header.html
share/apache/htdocs/manual/images/apache_pb.gif
share/apache/htdocs/manual/images/custom_errordocs.gif
share/apache/htdocs/manual/images/home.gif
share/apache/htdocs/manual/images/index.gif
share/apache/htdocs/manual/images/mod_rewrite_fig1.fig
share/apache/htdocs/manual/images/mod_rewrite_fig1.gif
share/apache/htdocs/manual/images/mod_rewrite_fig2.fig
share/apache/htdocs/manual/images/mod_rewrite_fig2.gif
share/apache/htdocs/manual/images/mod_ssl_sb.gif
share/apache/htdocs/manual/images/openssl_ics.gif
share/apache/htdocs/manual/images/sub.gif
share/apache/htdocs/manual/index.html
share/apache/htdocs/manual/install.html
share/apache/htdocs/manual/invoking.html
share/apache/htdocs/manual/keepalive.html
share/apache/htdocs/manual/location.html
share/apache/htdocs/manual/man-template.html
share/apache/htdocs/manual/install-tpf.html
share/apache/htdocs/manual/readme-tpf.html
share/apache/htdocs/manual/misc/API.html
share/apache/htdocs/manual/misc/FAQ.html
share/apache/htdocs/manual/misc/HTTP_Features.tsv
share/apache/htdocs/manual/misc/client_block_api.html
share/apache/htdocs/manual/misc/compat_notes.html
share/apache/htdocs/manual/misc/custom_errordocs.html
share/apache/htdocs/manual/misc/descriptors.html
share/apache/htdocs/manual/misc/fin_wait_2.html
share/apache/htdocs/manual/misc/footer.html
share/apache/htdocs/manual/misc/header.html
share/apache/htdocs/manual/misc/howto.html
share/apache/htdocs/manual/misc/index.html
share/apache/htdocs/manual/misc/known_client_problems.html
share/apache/htdocs/manual/misc/nopgp.html
share/apache/htdocs/manual/misc/perf-bsd44.html
share/apache/htdocs/manual/misc/perf-dec.html
share/apache/htdocs/manual/misc/perf-hp.html
share/apache/htdocs/manual/misc/perf-tuning.html
share/apache/htdocs/manual/misc/perf.html
share/apache/htdocs/manual/misc/security_tips.html
share/apache/htdocs/manual/misc/vif-info.html
share/apache/htdocs/manual/misc/windoz_keepalive.html
share/apache/htdocs/manual/mod/core.html
share/apache/htdocs/manual/mod/directive-dict.html
share/apache/htdocs/manual/mod/directives.html
share/apache/htdocs/manual/mod/footer.html
share/apache/htdocs/manual/mod/header.html
share/apache/htdocs/manual/mod/index.html
share/apache/htdocs/manual/mod/mod_define.html
share/apache/htdocs/manual/mod/mod_access.html
share/apache/htdocs/manual/mod/mod_actions.html
share/apache/htdocs/manual/mod/mod_alias.html
share/apache/htdocs/manual/mod/mod_asis.html
share/apache/htdocs/manual/mod/mod_auth.html
share/apache/htdocs/manual/mod/mod_auth_anon.html
share/apache/htdocs/manual/mod/mod_auth_db.html
share/apache/htdocs/manual/mod/mod_auth_dbm.html
share/apache/htdocs/manual/mod/mod_autoindex.html
share/apache/htdocs/manual/mod/mod_browser.html
share/apache/htdocs/manual/mod/mod_cern_meta.html
share/apache/htdocs/manual/mod/mod_cgi.html
share/apache/htdocs/manual/mod/mod_cookies.html
share/apache/htdocs/manual/mod/mod_digest.html
share/apache/htdocs/manual/mod/mod_dir.html
share/apache/htdocs/manual/mod/mod_dld.html
share/apache/htdocs/manual/mod/mod_dll.html
share/apache/htdocs/manual/mod/mod_env.html
share/apache/htdocs/manual/mod/mod_example.html
share/apache/htdocs/manual/mod/mod_expires.html
share/apache/htdocs/manual/mod/mod_headers.html
share/apache/htdocs/manual/mod/mod_imap.html
share/apache/htdocs/manual/mod/mod_include.html
share/apache/htdocs/manual/mod/mod_info.html
share/apache/htdocs/manual/mod/mod_isapi.html
share/apache/htdocs/manual/mod/mod_log_agent.html
share/apache/htdocs/manual/mod/mod_log_common.html
share/apache/htdocs/manual/mod/mod_log_config.html
share/apache/htdocs/manual/mod/mod_log_referer.html
share/apache/htdocs/manual/mod/mod_mime.html
share/apache/htdocs/manual/mod/mod_mime_magic.html
share/apache/htdocs/manual/mod/mod_mmap_static.html
share/apache/htdocs/manual/mod/mod_negotiation.html
share/apache/htdocs/manual/mod/mod_proxy.html
share/apache/htdocs/manual/mod/mod_rewrite.html
share/apache/htdocs/manual/mod/mod_setenvif.html
share/apache/htdocs/manual/mod/mod_so.html
share/apache/htdocs/manual/mod/mod_speling.html
share/apache/htdocs/manual/mod/mod_status.html
share/apache/htdocs/manual/mod/mod_unique_id.html
share/apache/htdocs/manual/mod/mod_userdir.html
share/apache/htdocs/manual/mod/mod_usertrack.html
share/apache/htdocs/manual/mod/mod_auth_digest.html
share/apache/htdocs/manual/mod/mod_vhost_alias.html
share/apache/htdocs/manual/mod/mod_ssl/index.html
share/apache/htdocs/manual/mod/mod_ssl/ssl_compat.gfont000.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_compat.html
share/apache/htdocs/manual/mod/mod_ssl/ssl_cover_logo.jpg
share/apache/htdocs/manual/mod/mod_ssl/ssl_cover_title.jpg
share/apache/htdocs/manual/mod/mod_ssl/ssl_faq.gfont000.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_faq.html
share/apache/htdocs/manual/mod/mod_ssl/ssl_glossary.html
share/apache/htdocs/manual/mod/mod_ssl/ssl_howto.gfont000.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_howto.html
share/apache/htdocs/manual/mod/mod_ssl/ssl_intro.gfont000.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_intro.html
share/apache/htdocs/manual/mod/mod_ssl/ssl_intro_fig1.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_intro_fig2.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_intro_fig3.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_overview.gfont000.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_overview.html
share/apache/htdocs/manual/mod/mod_ssl/ssl_overview_fig1.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_reference.gfont000.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_reference.html
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.head-chapter.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.head-num-1.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.head-num-2.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.head-num-3.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.head-num-4.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.head-num-5.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.head-num-6.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.head-num-7.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.imgdot-1x1-000000.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.imgdot-1x1-transp.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.navbut-next-n.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.navbut-next-s.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.navbut-prev-n.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.navbut-prev-s.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.title-abstract.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.title-compat.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.title-faq.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.title-gloss.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.title-howto.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.title-intro.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.title-over.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.title-preface.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.title-ref.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.title-toc.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.title-tutor.gif
share/apache/htdocs/manual/mod/mod_ssl/ssl_template.inc
share/apache/htdocs/manual/mod/mod_ssl/ssl_compat.wml
share/apache/htdocs/manual/mod/mod_ssl/ssl_cover.wml
share/apache/htdocs/manual/mod/mod_ssl/ssl_faq.wml
share/apache/htdocs/manual/mod/mod_ssl/ssl_glossary.wml
share/apache/htdocs/manual/mod/mod_ssl/ssl_howto.wml
share/apache/htdocs/manual/mod/mod_ssl/ssl_intro.wml
share/apache/htdocs/manual/mod/mod_ssl/ssl_overview.wml
share/apache/htdocs/manual/mod/mod_ssl/ssl_reference.wml
share/apache/htdocs/manual/multilogs.html
share/apache/htdocs/manual/new_features_1_0.html
share/apache/htdocs/manual/new_features_1_1.html
share/apache/htdocs/manual/new_features_1_2.html
share/apache/htdocs/manual/new_features_1_3.html
share/apache/htdocs/manual/process-model.html
share/apache/htdocs/manual/search/manual-index.cgi
share/apache/htdocs/manual/sections.html
share/apache/htdocs/manual/sourcereorg.html
share/apache/htdocs/manual/stopping.html
share/apache/htdocs/manual/suexec.html
share/apache/htdocs/manual/suexec_1_2.html
share/apache/htdocs/manual/unixware.html
share/apache/htdocs/manual/upgrading_to_1_3.html
share/apache/htdocs/manual/vhosts/details.html
share/apache/htdocs/manual/vhosts/details_1_2.html
share/apache/htdocs/manual/vhosts/examples.html
share/apache/htdocs/manual/vhosts/fd-limits.html
share/apache/htdocs/manual/vhosts/footer.html
share/apache/htdocs/manual/vhosts/header.html
share/apache/htdocs/manual/vhosts/host.html
share/apache/htdocs/manual/vhosts/index.html
share/apache/htdocs/manual/vhosts/ip-based.html
share/apache/htdocs/manual/vhosts/mass.html
share/apache/htdocs/manual/vhosts/name-based.html
share/apache/htdocs/manual/vhosts/vhosts-in-depth.html
share/apache/htdocs/manual/vhosts/virtual-host.html
share/apache/htdocs/manual/windows.html
share/apache/icons/README
share/apache/icons/a.gif
share/apache/icons/alert.black.gif
share/apache/icons/alert.red.gif
share/apache/icons/apache_pb.gif
share/apache/icons/back.gif
share/apache/icons/ball.gray.gif
share/apache/icons/ball.red.gif
share/apache/icons/binary.gif
share/apache/icons/binhex.gif
share/apache/icons/blank.gif
share/apache/icons/bomb.gif
share/apache/icons/box1.gif
share/apache/icons/box2.gif
share/apache/icons/broken.gif
share/apache/icons/burst.gif
share/apache/icons/c.gif
share/apache/icons/comp.blue.gif
share/apache/icons/comp.gray.gif
share/apache/icons/compressed.gif
share/apache/icons/continued.gif
share/apache/icons/dir.gif
share/apache/icons/down.gif
share/apache/icons/dvi.gif
share/apache/icons/f.gif
share/apache/icons/folder.gif
share/apache/icons/folder.open.gif
share/apache/icons/folder.sec.gif
share/apache/icons/forward.gif
share/apache/icons/generic.gif
share/apache/icons/generic.red.gif
share/apache/icons/generic.sec.gif
share/apache/icons/hand.right.gif
share/apache/icons/hand.up.gif
share/apache/icons/icon.sheet.gif
share/apache/icons/image1.gif
share/apache/icons/image2.gif
share/apache/icons/image3.gif
share/apache/icons/index.gif
share/apache/icons/layout.gif
share/apache/icons/left.gif
share/apache/icons/link.gif
share/apache/icons/movie.gif
share/apache/icons/p.gif
share/apache/icons/patch.gif
share/apache/icons/pdf.gif
share/apache/icons/pie0.gif
share/apache/icons/pie1.gif
share/apache/icons/pie2.gif
share/apache/icons/pie3.gif
share/apache/icons/pie4.gif
share/apache/icons/pie5.gif
share/apache/icons/pie6.gif
share/apache/icons/pie7.gif
share/apache/icons/pie8.gif
share/apache/icons/portal.gif
share/apache/icons/ps.gif
share/apache/icons/quill.gif
share/apache/icons/right.gif
share/apache/icons/screw1.gif
share/apache/icons/screw2.gif
share/apache/icons/script.gif
share/apache/icons/sound1.gif
share/apache/icons/sound2.gif
share/apache/icons/sphere1.gif
share/apache/icons/sphere2.gif
share/apache/icons/tar.gif
share/apache/icons/tex.gif
share/apache/icons/text.gif
share/apache/icons/transfer.gif
share/apache/icons/unknown.gif
share/apache/icons/up.gif
share/apache/icons/uu.gif
share/apache/icons/uuencoded.gif
share/apache/icons/world1.gif
share/apache/icons/world2.gif
share/apache/icons/small/README.txt
share/apache/icons/small/back.gif
share/apache/icons/small/binary.gif
share/apache/icons/small/binhex.gif
share/apache/icons/small/blank.gif
share/apache/icons/small/broken.gif
share/apache/icons/small/burst.gif
share/apache/icons/small/comp1.gif
share/apache/icons/small/comp2.gif
share/apache/icons/small/compressed.gif
share/apache/icons/small/continued.gif
share/apache/icons/small/dir.gif
share/apache/icons/small/dir2.gif
share/apache/icons/small/doc.gif
share/apache/icons/small/forward.gif
share/apache/icons/small/generic.gif
share/apache/icons/small/generic2.gif
share/apache/icons/small/generic3.gif
share/apache/icons/small/image.gif
share/apache/icons/small/image2.gif
share/apache/icons/small/index.gif
share/apache/icons/small/key.gif
share/apache/icons/small/movie.gif
share/apache/icons/small/patch.gif
share/apache/icons/small/ps.gif
share/apache/icons/small/rainbow.gif
share/apache/icons/small/sound.gif
share/apache/icons/small/sound2.gif
share/apache/icons/small/tar.gif
share/apache/icons/small/text.gif
share/apache/icons/small/transfer.gif
share/apache/icons/small/unknown.gif
share/apache/icons/small/uu.gif
@exec cd %D/etc/apache/ssl.crt && make >/dev/null 2>&1
@unexec rm -f %D/etc/apache/ssl.crt/[0-9a-fA-F]*.[0-9]*
@dirrm etc/apache/ssl.crt
@dirrm etc/apache/ssl.crl
@dirrm etc/apache/ssl.csr
@dirrm etc/apache/ssl.prm
@dirrm etc/apache/ssl.key
@dirrm etc/apache
@dirrm include/apache/xml/
@dirrm include/apache
@dirrm libexec/apache
@dirrm share/apache/cgi-bin
@dirrm share/apache/htdocs/manual/images
@dirrm share/apache/htdocs/manual/misc
@dirrm share/apache/htdocs/manual/mod/mod_ssl
@dirrm share/apache/htdocs/manual/mod
@dirrm share/apache/htdocs/manual/search
@dirrm share/apache/htdocs/manual/vhosts
@dirrm share/apache/htdocs/manual
@dirrm share/apache/htdocs
@dirrm share/apache/icons/small
@dirrm share/apache/icons
@dirrm share/apache
