# www/apache13-modssl/pkg-plist
bin/dbmmanage
bin/htdigest
bin/htpasswd
etc/apache/access.conf
etc/apache/access.conf.default
etc/apache/httpd.conf
etc/apache/httpd.conf.default
etc/apache/magic
etc/apache/magic.default
etc/apache/mime.types
etc/apache/mime.types.default
etc/apache/srm.conf
etc/apache/srm.conf.default
etc/apache/ssl.crt/README.CRT
etc/apache/ssl.crt/Makefile
etc/apache/ssl.crt/ca-bundle.crt
etc/apache/ssl.crt/snakeoil-ca.crt
etc/apache/ssl.crt/server.crt
etc/apache/ssl.csr/README.CSR
etc/apache/ssl.csr/server.csr
etc/apache/ssl.key/README.KEY
etc/apache/ssl.key/snakeoil-ca.key
etc/apache/ssl.key/server.key
etc/rc.d/apache.sh
include/apache/alloc.h
include/apache/ap_ctx.h
include/apache/ap_hook.h
include/apache/ap.h
include/apache/ap_compat.h
include/apache/ap_config.h
include/apache/ap_config_auto.h
include/apache/ap_ctype.h
include/apache/ap_md5.h
include/apache/ap_mmn.h
include/apache/buff.h
include/apache/compat.h
include/apache/conf.h
include/apache/explain.h
include/apache/fnmatch.h
include/apache/hsregex.h
include/apache/http_conf_globals.h
include/apache/http_config.h
include/apache/http_core.h
include/apache/http_log.h
include/apache/http_main.h
include/apache/http_protocol.h
include/apache/http_request.h
include/apache/http_vhost.h
include/apache/httpd.h
include/apache/multithread.h
include/apache/os.h
include/apache/os-inline.c
include/apache/rfc1413.h
include/apache/scoreboard.h
include/apache/util_date.h
include/apache/util_md5.h
include/apache/util_script.h
include/apache/util_uri.h
libexec/apache/libproxy.so
libexec/apache/libssl.so
libexec/apache/mod_access.so
libexec/apache/mod_actions.so
libexec/apache/mod_alias.so
libexec/apache/mod_asis.so
libexec/apache/mod_auth.so
libexec/apache/mod_auth_anon.so
libexec/apache/mod_auth_db.so
libexec/apache/mod_autoindex.so
libexec/apache/mod_cern_meta.so
libexec/apache/mod_cgi.so
libexec/apache/mod_digest.so
libexec/apache/mod_dir.so
libexec/apache/mod_env.so
libexec/apache/mod_expires.so
libexec/apache/mod_headers.so
libexec/apache/mod_imap.so
libexec/apache/mod_include.so
libexec/apache/mod_info.so
libexec/apache/mod_log_config.so
libexec/apache/mod_mime.so
libexec/apache/mod_mime_magic.so
libexec/apache/mod_negotiation.so
libexec/apache/mod_rewrite.so
libexec/apache/mod_setenvif.so
libexec/apache/mod_speling.so
libexec/apache/mod_status.so
libexec/apache/mod_unique_id.so
libexec/apache/mod_userdir.so
libexec/apache/mod_usertrack.so
sbin/ab
sbin/apachectl
sbin/apxs
sbin/httpd
sbin/logresolve
sbin/rotatelogs
share/doc/apache/apache_pb.gif
share/doc/apache/index.html
share/doc/apache/manual/LICENSE
share/doc/apache/manual/bind.html
share/doc/apache/manual/cgi_path.html
share/doc/apache/manual/content-negotiation.html
share/doc/apache/manual/custom-error.html
share/doc/apache/manual/dns-caveats.html
share/doc/apache/manual/dso.html
share/doc/apache/manual/ebcdic.html
share/doc/apache/manual/env.html
share/doc/apache/manual/footer.html
share/doc/apache/manual/handler.html
share/doc/apache/manual/header.html
share/doc/apache/manual/images/apache_pb.gif
share/doc/apache/manual/images/custom_errordocs.gif
share/doc/apache/manual/images/home.gif
share/doc/apache/manual/images/index.gif
share/doc/apache/manual/images/mod_rewrite_fig1.fig
share/doc/apache/manual/images/mod_rewrite_fig1.gif
share/doc/apache/manual/images/mod_rewrite_fig2.fig
share/doc/apache/manual/images/mod_rewrite_fig2.gif
share/doc/apache/manual/images/mod_ssl_sb.gif
share/doc/apache/manual/images/ssleay.gif
share/doc/apache/manual/images/sub.gif
share/doc/apache/manual/index.html
share/doc/apache/manual/install.html
share/doc/apache/manual/invoking.html
share/doc/apache/manual/keepalive.html
share/doc/apache/manual/location.html
share/doc/apache/manual/man-template.html
share/doc/apache/manual/install-tpf.html
share/doc/apache/manual/readme-tpf.html
share/doc/apache/manual/misc/API.html
share/doc/apache/manual/misc/FAQ.html
share/doc/apache/manual/misc/HTTP_Features.tsv
share/doc/apache/manual/misc/client_block_api.html
share/doc/apache/manual/misc/compat_notes.html
share/doc/apache/manual/misc/custom_errordocs.html
share/doc/apache/manual/misc/descriptors.html
share/doc/apache/manual/misc/fin_wait_2.html
share/doc/apache/manual/misc/footer.html
share/doc/apache/manual/misc/header.html
share/doc/apache/manual/misc/howto.html
share/doc/apache/manual/misc/index.html
share/doc/apache/manual/misc/known_client_problems.html
share/doc/apache/manual/misc/nopgp.html
share/doc/apache/manual/misc/perf-bsd44.html
share/doc/apache/manual/misc/perf-dec.html
share/doc/apache/manual/misc/perf-hp.html
share/doc/apache/manual/misc/perf-tuning.html
share/doc/apache/manual/misc/perf.html
share/doc/apache/manual/misc/security_tips.html
share/doc/apache/manual/misc/vif-info.html
share/doc/apache/manual/misc/windoz_keepalive.html
share/doc/apache/manual/mod/core.html
share/doc/apache/manual/mod/directive-dict.html
share/doc/apache/manual/mod/directives.html
share/doc/apache/manual/mod/footer.html
share/doc/apache/manual/mod/header.html
share/doc/apache/manual/mod/index.html
share/doc/apache/manual/mod/mod_access.html
share/doc/apache/manual/mod/mod_actions.html
share/doc/apache/manual/mod/mod_alias.html
share/doc/apache/manual/mod/mod_asis.html
share/doc/apache/manual/mod/mod_auth.html
share/doc/apache/manual/mod/mod_auth_anon.html
share/doc/apache/manual/mod/mod_auth_db.html
share/doc/apache/manual/mod/mod_auth_dbm.html
share/doc/apache/manual/mod/mod_autoindex.html
share/doc/apache/manual/mod/mod_browser.html
share/doc/apache/manual/mod/mod_cern_meta.html
share/doc/apache/manual/mod/mod_cgi.html
share/doc/apache/manual/mod/mod_cookies.html
share/doc/apache/manual/mod/mod_digest.html
share/doc/apache/manual/mod/mod_dir.html
share/doc/apache/manual/mod/mod_dld.html
share/doc/apache/manual/mod/mod_dll.html
share/doc/apache/manual/mod/mod_env.html
share/doc/apache/manual/mod/mod_example.html
share/doc/apache/manual/mod/mod_expires.html
share/doc/apache/manual/mod/mod_headers.html
share/doc/apache/manual/mod/mod_imap.html
share/doc/apache/manual/mod/mod_include.html
share/doc/apache/manual/mod/mod_info.html
share/doc/apache/manual/mod/mod_isapi.html
share/doc/apache/manual/mod/mod_log_agent.html
share/doc/apache/manual/mod/mod_log_common.html
share/doc/apache/manual/mod/mod_log_config.html
share/doc/apache/manual/mod/mod_log_referer.html
share/doc/apache/manual/mod/mod_mime.html
share/doc/apache/manual/mod/mod_mime_magic.html
share/doc/apache/manual/mod/mod_mmap_static.html
share/doc/apache/manual/mod/mod_negotiation.html
share/doc/apache/manual/mod/mod_proxy.html
share/doc/apache/manual/mod/mod_rewrite.html
share/doc/apache/manual/mod/mod_setenvif.html
share/doc/apache/manual/mod/mod_so.html
share/doc/apache/manual/mod/mod_speling.html
share/doc/apache/manual/mod/mod_status.html
share/doc/apache/manual/mod/mod_unique_id.html
share/doc/apache/manual/mod/mod_userdir.html
share/doc/apache/manual/mod/mod_usertrack.html
share/doc/apache/manual/mod/mod_ssl/index.html
share/doc/apache/manual/mod/mod_ssl/ssl_compat.gfont000.gif
share/doc/apache/manual/mod/mod_ssl/ssl_compat.html
share/doc/apache/manual/mod/mod_ssl/ssl_cover_logo.jpg
share/doc/apache/manual/mod/mod_ssl/ssl_cover_title.gif
share/doc/apache/manual/mod/mod_ssl/ssl_faq.gfont000.gif
share/doc/apache/manual/mod/mod_ssl/ssl_faq.html
share/doc/apache/manual/mod/mod_ssl/ssl_glossary.html
share/doc/apache/manual/mod/mod_ssl/ssl_howto.gfont000.gif
share/doc/apache/manual/mod/mod_ssl/ssl_howto.html
share/doc/apache/manual/mod/mod_ssl/ssl_intro.gfont000.gif
share/doc/apache/manual/mod/mod_ssl/ssl_intro.html
share/doc/apache/manual/mod/mod_ssl/ssl_intro_fig1.gif
share/doc/apache/manual/mod/mod_ssl/ssl_intro_fig2.gif
share/doc/apache/manual/mod/mod_ssl/ssl_intro_fig3.gif
share/doc/apache/manual/mod/mod_ssl/ssl_overview.gfont000.gif
share/doc/apache/manual/mod/mod_ssl/ssl_overview.html
share/doc/apache/manual/mod/mod_ssl/ssl_overview_fig1.gif
share/doc/apache/manual/mod/mod_ssl/ssl_reference.gfont000.gif
share/doc/apache/manual/mod/mod_ssl/ssl_reference.html
share/doc/apache/manual/mod/mod_ssl/ssl_template.head-chapter.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.head-num-1.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.head-num-2.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.head-num-3.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.head-num-4.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.head-num-5.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.head-num-6.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.head-num-7.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.imgdot-1x1-000000.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.imgdot-1x1-transp.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.navbut-next-n.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.navbut-next-s.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.navbut-prev-n.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.navbut-prev-s.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.title-abstract.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.title-compat.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.title-faq.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.title-gloss.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.title-howto.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.title-intro.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.title-over.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.title-preface.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.title-ref.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.title-toc.gif
share/doc/apache/manual/mod/mod_ssl/ssl_template.title-tutor.gif
share/doc/apache/manual/mod/mod_ssl/ssl_compat.wml
share/doc/apache/manual/mod/mod_ssl/ssl_cover.wml
share/doc/apache/manual/mod/mod_ssl/ssl_faq.wml
share/doc/apache/manual/mod/mod_ssl/ssl_glossary.wml
share/doc/apache/manual/mod/mod_ssl/ssl_howto.wml
share/doc/apache/manual/mod/mod_ssl/ssl_intro.wml
share/doc/apache/manual/mod/mod_ssl/ssl_overview.wml
share/doc/apache/manual/mod/mod_ssl/ssl_reference.wml
share/doc/apache/manual/multilogs.html
share/doc/apache/manual/new_features_1_0.html
share/doc/apache/manual/new_features_1_1.html
share/doc/apache/manual/new_features_1_2.html
share/doc/apache/manual/new_features_1_3.html
share/doc/apache/manual/process-model.html
share/doc/apache/manual/sections.html
share/doc/apache/manual/sourcereorg.html
share/doc/apache/manual/stopping.html
share/doc/apache/manual/suexec.html
share/doc/apache/manual/suexec_1_2.html
share/doc/apache/manual/unixware.html
share/doc/apache/manual/upgrading_to_1_3.html
share/doc/apache/manual/vhosts/details.html
share/doc/apache/manual/vhosts/details_1_2.html
share/doc/apache/manual/vhosts/examples.html
share/doc/apache/manual/vhosts/fd-limits.html
share/doc/apache/manual/vhosts/footer.html
share/doc/apache/manual/vhosts/header.html
share/doc/apache/manual/vhosts/host.html
share/doc/apache/manual/vhosts/index.html
share/doc/apache/manual/vhosts/ip-based.html
share/doc/apache/manual/vhosts/name-based.html
share/doc/apache/manual/vhosts/vhosts-in-depth.html
share/doc/apache/manual/vhosts/virtual-host.html
share/doc/apache/manual/windows.html
www/cgi-bin.default/printenv
www/cgi-bin.default/test-cgi
www/icons/README
www/icons/a.gif
www/icons/alert.black.gif
www/icons/alert.red.gif
www/icons/apache_pb.gif
www/icons/back.gif
www/icons/ball.gray.gif
www/icons/ball.red.gif
www/icons/binary.gif
www/icons/binhex.gif
www/icons/blank.gif
www/icons/bomb.gif
www/icons/box1.gif
www/icons/box2.gif
www/icons/broken.gif
www/icons/burst.gif
www/icons/c.gif
www/icons/comp.blue.gif
www/icons/comp.gray.gif
www/icons/compressed.gif
www/icons/continued.gif
www/icons/dir.gif
www/icons/down.gif
www/icons/dvi.gif
www/icons/f.gif
www/icons/folder.gif
www/icons/folder.open.gif
www/icons/folder.sec.gif
www/icons/forward.gif
www/icons/generic.gif
www/icons/generic.red.gif
www/icons/generic.sec.gif
www/icons/hand.right.gif
www/icons/hand.up.gif
www/icons/icon.sheet.gif
www/icons/image1.gif
www/icons/image2.gif
www/icons/image3.gif
www/icons/index.gif
www/icons/layout.gif
www/icons/left.gif
www/icons/link.gif
www/icons/movie.gif
www/icons/p.gif
www/icons/patch.gif
www/icons/pdf.gif
www/icons/pie0.gif
www/icons/pie1.gif
www/icons/pie2.gif
www/icons/pie3.gif
www/icons/pie4.gif
www/icons/pie5.gif
www/icons/pie6.gif
www/icons/pie7.gif
www/icons/pie8.gif
www/icons/portal.gif
www/icons/ps.gif
www/icons/quill.gif
www/icons/right.gif
www/icons/screw1.gif
www/icons/screw2.gif
www/icons/script.gif
www/icons/sound1.gif
www/icons/sound2.gif
www/icons/sphere1.gif
www/icons/sphere2.gif
www/icons/tar.gif
www/icons/tex.gif
www/icons/text.gif
www/icons/transfer.gif
www/icons/unknown.gif
www/icons/up.gif
www/icons/uu.gif
www/icons/uuencoded.gif
www/icons/world1.gif
www/icons/world2.gif
www/icons/small/README.txt
www/icons/small/back.gif
www/icons/small/binary.gif
www/icons/small/binhex.gif
www/icons/small/blank.gif
www/icons/small/broken.gif
www/icons/small/burst.gif
www/icons/small/comp1.gif
www/icons/small/comp2.gif
www/icons/small/compressed.gif
www/icons/small/continued.gif
www/icons/small/dir.gif
www/icons/small/dir2.gif
www/icons/small/forward.gif
www/icons/small/generic.gif
www/icons/small/generic2.gif
www/icons/small/generic3.gif
www/icons/small/image.gif
www/icons/small/image2.gif
www/icons/small/index.gif
www/icons/small/key.gif
www/icons/small/movie.gif
www/icons/small/rainbow.gif
www/icons/small/sound.gif
www/icons/small/sound2.gif
www/icons/small/tar.gif
www/icons/small/text.gif
www/icons/small/transfer.gif
www/icons/small/unknown.gif
www/icons/small/uu.gif
@exec cd %D/etc/apache/ssl.crt && make >/dev/null 2>&1
@exec mkdir -p %D/www
@exec mkdir -p %D/www/proxy
@exec ln -fs %D/www/data.default %D/www/data
@exec ln -fs %D/www/cgi-bin.default %D/www/cgi-bin
@unexec rm -f %D/www/data
@unexec rm -f %D/www/data.default
@unexec rm -f %D/www/cgi-bin
@dirrm www/cgi-bin.default
@dirrm www/proxy
@dirrm www/icons/small
@dirrm www/icons
@dirrm www
@unexec rm -f %D/etc/apache/ssl.crt/[0-9a-fA-F]*.[0-9]*
@dirrm etc/apache/ssl.crt
@dirrm etc/apache/ssl.csr
@dirrm etc/apache/ssl.key
@dirrm etc/apache
@dirrm include/apache
@dirrm libexec/apache
@dirrm share/doc/apache/manual/images
@dirrm share/doc/apache/manual/misc
@dirrm share/doc/apache/manual/mod/mod_ssl
@dirrm share/doc/apache/manual/mod
@dirrm share/doc/apache/manual/vhosts
@dirrm share/doc/apache/manual
@dirrm share/doc/apache
