# dns/gdnsd2/files/gdnsd.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: gdnsd
# REQUIRE: DAEMON
# REQUIRE: LOGIN
# KEYWORD: FreeBSD shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# gdnsd_enable (bool):	Set to NO by default.
#				Set it to YES to enable gdnsd.
# gdnsd_flags (str):		Set to "" by default
#				Extra command line argument flags

. /etc/rc.subr

name="gdnsd"
rcvar=gdnsd_enable

# set defaults
load_rc_config $name
: ${gdnsd_enable:="NO"}
: ${gdnsd_flags:=""}

command="%%PREFIX%%/sbin/gdnsd"
command_args="start"
procname=${command}
reload_cmd="gdnsd_reload"
configtest_cmd="gdnsd_checkconf"
extra_commands="reload configtest"
gdnsd_reload()
{
  ${command} ${gdnsd_flags} reload-zones
}

gdnsd_checkconf()
{
  ${command} ${gdnsd_flags} checkconf
}

run_rc_command "$1"

