# sysutils/py-salt/files/openrc-salt_master.in
#!/sbin/openrc-run

# Configuration options for this service
# =======================
# salt_master_paths (string):      Set to "/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin:/usr/local/sbin" by default.
#               Default $PATH for salt_master
# salt_master_eggcache (string):   Set to "/tmp" by default.
#               Allows defining egg cache directory to fix runtime on diskless systems.

name=salt_master
command="/usr/local/bin/salt-master"
required_dirs="/usr/local/etc/salt"
command_args="-c ${required_dirs}"
pidfile="${salt_master_pidfile:-/var/run/salt-master.pid}"
command_background="true"

start_stop_daemon_args=" -e PATH=${salt_master_paths:-/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin:/usr/local/sbin} -e PYTHON_EGG_CACHE=\"${salt_master_eggcache:-/tmp}\""

depend() {
        keyword -shutdown
}
