# sysutils/py-salt/files/salt_master.in
#!/bin/sh

# $FreeBSD$
#
# Salt Master startup script
#
# PROVIDE: salt_master
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following to /etc/rc.conf[.local] to enable this service
#
# salt_master_enable="YES"
#

. /etc/rc.subr

name="salt_master"
rcvar=salt_master_enable

load_rc_config ${name}

: ${salt_master_enable:="NO"}

command="%%PREFIX%%/bin/salt-master"
command_interpreter="%%PYTHON_CMD%%"
required_files="%%PREFIX%%/etc/salt"
command_args="-c ${required_files} -d"

run_rc_command "$1"
