# databases/pgpool-II-33/files/openrc-pgpool.in
#!/sbin/openrc-run
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# pgpool_conf (file):		Set location of your config.
#				Default is "%%PREFIX%%/etc/pgpool.conf"
# pgpool_flags (string):	Default is "-f $pgpool_conf"
# pgpool_looger (bool):		Use logger for logging.
# pgpool_looger_pri (string):	logger priority (default "local0.info").
# pgpool_user (uid):		User ID to run as (default nobody)

name=pgpool

: ${pgpool_conf="%%PREFIX%%/etc/pgpool.conf"}
: ${pgpool_flags="-f ${pgpool_conf}"}
: ${pgpool_user="nobody"}
: ${pgpool_logger="YES"}
: ${pgpool_logger_pri="local0.info"}

command="%%PREFIX%%/bin/pgpool"
command_args="${pgpool_flags}"
command_user=${pgpool_user}
required_files="${pgpool_conf}"

start_pre()
{
    if yesno pgpool_logger; then
		command_args="$command_args -n 2>&1 | logger -t $name -p $pgpool_logger_pri &"
	fi
}

stop_pre() {
	$command $command_args stop
}

depend() {
	keyword -shutdown
}
