# www/apache13/pkg-plist
etc/rc.d/apache.sh
@unexec rm %D/etc/rc.d/apache.sh
sbin/httpd
www/cgi-bin
www/data
www/server/logs
www/server/conf/httpd.conf
www/server/conf/srm.conf
www/server/conf/access.conf
www/server/conf/mime.types
www/icons/README
www/icons/a.gif
www/icons/alert.black.gif
www/icons/alert.red.gif
www/icons/apache_pb.gif
www/icons/back.gif
www/icons/ball.gray.gif
www/icons/ball.red.gif
www/icons/binary.gif
www/icons/binhex.gif
www/icons/blank.gif
www/icons/bomb.gif
www/icons/box1.gif
www/icons/box2.gif
www/icons/broken.gif
www/icons/burst.gif
www/icons/c.gif
www/icons/comp.blue.gif
www/icons/comp.gray.gif
www/icons/compressed.gif
www/icons/continued.gif
www/icons/dir.gif
www/icons/down.gif
www/icons/dvi.gif
www/icons/f.gif
www/icons/folder.gif
www/icons/folder.open.gif
www/icons/folder.sec.gif
www/icons/forward.gif
www/icons/generic.gif
www/icons/generic.red.gif
www/icons/generic.sec.gif
www/icons/hand.right.gif
www/icons/hand.up.gif
www/icons/icon.sheet.gif
www/icons/image1.gif
www/icons/image2.gif
www/icons/image3.gif
www/icons/index.gif
www/icons/layout.gif
www/icons/left.gif
www/icons/link.gif
www/icons/movie.gif
www/icons/p.gif
www/icons/patch.gif
www/icons/pdf.gif
www/icons/pie0.gif
www/icons/pie1.gif
www/icons/pie2.gif
www/icons/pie3.gif
www/icons/pie4.gif
www/icons/pie5.gif
www/icons/pie6.gif
www/icons/pie7.gif
www/icons/pie8.gif
www/icons/portal.gif
www/icons/ps.gif
www/icons/quill.gif
www/icons/right.gif
www/icons/screw1.gif
www/icons/screw2.gif
www/icons/script.gif
www/icons/sound1.gif
www/icons/sound2.gif
www/icons/sphere1.gif
www/icons/sphere2.gif
www/icons/tar.gif
www/icons/tex.gif
www/icons/text.gif
www/icons/transfer.gif
www/icons/unknown.gif
www/icons/up.gif
www/icons/uu.gif
www/icons/uuencoded.gif
www/icons/world1.gif
www/icons/world2.gif
