# japanese/ebnetd/pkg-plist
etc/ndtpd.conf.sample
etc/rc.d/ndtpd.sh.sample
@unexec install-info --delete %D/info/ndtpd-ja.info %D/info/dir
@unexec install-info --delete %D/info/ndtpd.info %D/info/dir
info/ndtpd-ja.info
info/ndtpd-ja.info-1
info/ndtpd-ja.info-2
info/ndtpd.info
info/ndtpd.info-1
info/ndtpd.info-2
@exec install-info %D/info/ndtpd-ja.info %D/info/dir
@exec install-info %D/info/ndtpd.info %D/info/dir
libexec/ndtpstat
sbin/ndtpcheck
sbin/ndtpcontrol
sbin/ndtpd
sbin/ndtpdaily
sbin/ndtpupgrade
