# net/netatalk-devel/pkg-plist
bin/aecho
bin/getzones
bin/nbplkup
bin/nbprgstr
bin/nbpunrgstr
bin/psorder
bin/megatron
bin/unbin
bin/unhex
bin/unsingle
bin/hqx2bin
bin/single2bin
bin/macbinary
bin/pap
bin/papstatus
etc/rc.d/netatalk.sh
etc/AppleVolumes.default
etc/AppleVolumes.system
include/atalk/adouble.h
include/atalk/rtmp.h
include/atalk/paths.h
include/atalk/pap.h
include/atalk/nbp.h
include/atalk/ddp.h
include/atalk/compat.h
include/atalk/atp.h
include/atalk/asp.h
include/atalk/afp.h
include/atalk/aep.h
include/atalk/zip.h
include/netatalk/ddp.h
include/netatalk/endian.h
include/netatalk/ddp_var.h
include/netatalk/phase2.h
include/netatalk/at_var.h
include/netatalk/at.h
include/netatalk/aarp.h
lib/libatalk.a
lib/libatalk_p.a
libexec/afpd
libexec/atalkd
libexec/papd
libexec/psa
libexec/etc2ps
libexec/psf
man/man1/aecho.1.gz
man/man1/unsingle.1.gz
man/man1/nbp.1.gz
man/man1/psorder.1.gz
man/man1/unbin.1.gz
man/man1/unhex.1.gz
man/man1/nbplkup.1.gz
man/man1/hqx2bin.1.gz
man/man1/pap.1.gz
man/man1/macbinary.1.gz
man/man1/getzones.1.gz
man/man1/nbprgstr.1.gz
man/man1/megatron.1.gz
man/man1/single2bin.1.gz
man/man1/papstatus.1.gz
man/man3/nbp_name.3.gz
man/man3/atalk_aton.3.gz
man/man4/atalk.4.gz
man/man8/papd.8.gz
man/man8/psf.8.gz
man/man8/afpd.8.gz
man/man8/atalkd.8.gz
share/examples/netatalk/atalkd.conf
share/examples/netatalk/papd.conf
share/filters/ofpap
share/filters/ifpap
share/filters/tfpap
share/filters/ifpaprev
share/filters/tfpaprev
share/filters/ofwpap
share/filters/ifwpap
share/filters/tfwpap
share/filters/ifwpaprev
share/filters/tfwpaprev
share/filters/ofmpap
share/filters/ifmpap
share/filters/tfmpap
share/filters/ifmpaprev
share/filters/tfmpaprev
share/filters/ofwmpap
share/filters/ifwmpap
share/filters/tfwmpap
share/filters/ifwmpaprev
share/filters/tfwmpaprev
share/pagecount.ps
