# net/netatalk-devel/pkg-plist
bin/aecho
bin/getzones
bin/hqx2bin
bin/macbinary
bin/megatron
bin/nbplkup
bin/nbprgstr
bin/nbpunrgstr
bin/pap
bin/papstatus
bin/psorder
bin/single2bin
bin/unbin
bin/unhex
bin/unsingle
etc/AppleVolumes.default.dist
etc/AppleVolumes.system.dist
etc/rc.d/netatalk.sh
include/atalk/adouble.h
include/atalk/aep.h
include/atalk/afp.h
include/atalk/asp.h
include/atalk/atp.h
include/atalk/compat.h
include/atalk/ddp.h
include/atalk/nbp.h
include/atalk/pap.h
include/atalk/paths.h
include/atalk/rtmp.h
include/atalk/zip.h
@dirrm include/atalk
lib/libatalk.a
lib/libatalk_p.a
libexec/afpd
libexec/atalkd
libexec/etc2ps
libexec/papd
libexec/psa
libexec/psf
share/examples/netatalk/atalkd.conf
share/examples/netatalk/papd.conf
@dirrm share/examples/netatalk
share/filters/ifmpap
share/filters/ifmpaprev
share/filters/ifpap
share/filters/ifpaprev
share/filters/ifwmpap
share/filters/ifwmpaprev
share/filters/ifwpap
share/filters/ifwpaprev
share/filters/ofmpap
share/filters/ofpap
share/filters/ofwmpap
share/filters/ofwpap
share/filters/tfmpap
share/filters/tfmpaprev
share/filters/tfpap
share/filters/tfpaprev
share/filters/tfwmpap
share/filters/tfwmpaprev
share/filters/tfwpap
share/filters/tfwpaprev
@dirrm share/filters
share/pagecount.ps
@comment	Remove man page aliases
@unexec rm -f %D/man/man1/hqx2bin.1*
@unexec rm -f %D/man/man1/macbinary.1*
@unexec rm -f %D/man/man1/single2bin.1*
@unexec rm -f %D/man/man1/unbin.1*
@unexec rm -f %D/man/man1/unhex.1*
@unexec rm -f %D/man/man1/unsingle.1*
@unexec rm -f %D/man/man1/nbplkup.1*
@unexec rm -f %D/man/man1/nbprgstr.1*
@unexec rm -f %D/man/man1/nbpunrgstr.1*
@unexec rm -f %D/man/man1/papstatus.1*
