# net/netatalk-devel/pkg-plist
bin/aecho
bin/getzones
bin/hqx2bin
bin/macbinary
bin/megatron
bin/nbplkup
bin/nbprgstr
bin/nbpunrgstr
bin/pap
bin/papstatus
bin/psorder
bin/single2bin
bin/unbin
bin/unhex
bin/unsingle
etc/AppleVolumes.default.dist
etc/AppleVolumes.system.dist
etc/rc.d/netatalk.sh
include/atalk/adouble.h
include/atalk/aep.h
include/atalk/afp.h
include/atalk/asp.h
include/atalk/atp.h
include/atalk/compat.h
include/atalk/ddp.h
include/atalk/nbp.h
include/atalk/pap.h
include/atalk/paths.h
include/atalk/rtmp.h
include/atalk/zip.h
@dirrm include/atalk
include/netatalk/aarp.h
include/netatalk/at.h
include/netatalk/at_var.h
include/netatalk/ddp.h
include/netatalk/ddp_var.h
include/netatalk/endian.h
include/netatalk/phase2.h
@dirrm include/netatalk
lib/libatalk.a
lib/libatalk_p.a
libexec/afpd
libexec/atalkd
libexec/etc2ps
libexec/papd
libexec/psa
libexec/psf
man/man1/aecho.1
man/man1/getzones.1
man/man1/hqx2bin.1
man/man1/macbinary.1
man/man1/megatron.1
man/man1/nbp.1
man/man1/nbplkup.1
man/man1/nbprgstr.1
man/man1/pap.1
man/man1/papstatus.1
man/man1/psorder.1
man/man1/single2bin.1
man/man1/unbin.1
man/man1/unhex.1
man/man1/unsingle.1
man/man3/atalk_aton.3
man/man3/nbp_name.3
man/man4/atalk.4
man/man8/afpd.8
man/man8/atalkd.8
man/man8/papd.8
man/man8/psf.8
share/examples/netatalk/atalkd.conf
share/examples/netatalk/papd.conf
@dirrm share/examples/netatalk
share/filters/ifmpap
share/filters/ifmpaprev
share/filters/ifpap
share/filters/ifpaprev
share/filters/ifwmpap
share/filters/ifwmpaprev
share/filters/ifwpap
share/filters/ifwpaprev
share/filters/ofmpap
share/filters/ofpap
share/filters/ofwmpap
share/filters/ofwpap
share/filters/tfmpap
share/filters/tfmpaprev
share/filters/tfpap
share/filters/tfpaprev
share/filters/tfwmpap
share/filters/tfwmpaprev
share/filters/tfwpap
share/filters/tfwpaprev
@dirrm share/filters
share/pagecount.ps
