# net/priss/pkg-plist
etc/QTSS.conf
etc/rc.d/qtss.sh
sbin/QuickTimeStreamingServer
@dirrm movies
