# ftp/proftpd-devel/pkg-plist
bin/ftpcount
bin/ftpwho
etc/proftpd.conf.default
etc/rc.d/proftpd.sh.sample
@exec if [ ! -f %D/etc/proftpd.conf ] ; then cp -p %D/%F %B/proftpd.conf; fi
libexec/proftpd
sbin/ftpshut
share/doc/proftpd/Configuration.html
share/doc/proftpd/FAQ-config.html
@dirrm share/doc/proftpd
