# net/netatalk/pkg-plist
bin/achfile
bin/add_netatalk_printer
bin/adv1tov2
bin/aecho
bin/afile
bin/afppasswd
bin/apple_cp
bin/apple_mv
bin/apple_rm
bin/cleanappledouble.pl
bin/getzones
bin/macusers
bin/makecode
bin/megatron
@unexec rm -f %B/unbin
@unexec rm -f %B/unhex
@unexec rm -f %B/unsingle
@unexec rm -f %B/hqx2bin
@unexec rm -f %B/single2bin
@unexec rm -f %B/macbinary
@unexec rm -f %B/binheader
@unexec rm -f %B/nadheader
@exec cd %B && ln -s megatron unbin
@exec cd %B && ln -s megatron unhex
@exec cd %B && ln -s megatron unsingle
@exec cd %B && ln -s megatron hqx2bin
@exec cd %B && ln -s megatron single2bin
@exec cd %B && ln -s megatron macbinary
@exec cd %B && ln -s megatron binheader
@exec cd %B && ln -s megatron nadheader
%%NETATALKCNID%%bin/cnid_didname_verify
bin/nbplkup
bin/nbprgstr
bin/nbpunrgstr
bin/netatalk-config
bin/netatalkshorternamelinks.pl
bin/pap
bin/papstatus
bin/parsecode
bin/psorder
bin/timeout
sbin/afpd
sbin/atalkd
sbin/ifmpap
sbin/ifmpaprev
sbin/ifpap
sbin/ifpaprev
sbin/ifwmpap
sbin/ifwmpaprev
sbin/ifwpap
sbin/ifwpaprev
sbin/ofmpap
sbin/ofpap
sbin/ofwmpap
sbin/ofwpap
sbin/papd
sbin/psa
sbin/psf
sbin/tfmpap
sbin/tfmpaprev
sbin/tfpap
sbin/tfpaprev
sbin/tfwmpap
sbin/tfwmpaprev
sbin/tfwpap
sbin/tfwpaprev
@unexec if cmp -s %D/etc/AppleVolumes.default %D/etc/AppleVolumes.default.dist; then rm -f %D/etc/AppleVolumes.default; fi
etc/AppleVolumes.default.dist
@exec [ ! -f %B/AppleVolumes.default ] && cp %B/%f %B/AppleVolumes.default
@unexec if cmp -s %D/etc/AppleVolumes.system %D/etc/AppleVolumes.system.dist; then rm -f %D/etc/AppleVolumes.system; fi
etc/AppleVolumes.system.dist
@exec [ ! -f %B/AppleVolumes.system ] && cp %B/%f %B/AppleVolumes.system
@unexec if cmp -s %D/etc/afpd.conf %D/etc/afpd.conf.dist; then rm -f %D/etc/afpd.conf; fi
etc/afpd.conf.dist
@exec [ ! -f %B/afpd.conf ] && cp %B/%f %B/afpd.conf
@unexec if cmp -s %D/etc/atalkd.conf %D/etc/atalkd.conf.dist; then rm -f %D/etc/atalkd.conf; fi
etc/atalkd.conf.dist
@exec [ ! -f %B/atalkd.conf ] && cp %B/%f %B/atalkd.conf
@unexec if cmp -s %D/etc/papd.conf %D/etc/papd.conf.dist; then rm -f %D/etc/papd.conf; fi
etc/papd.conf.dist
@exec [ ! -f %B/papd.conf ] && cp %B/%f %B/papd.conf
etc/nls/maccode.437
etc/nls/maccode.850
etc/nls/maccode.iso8859-1
etc/nls/maccode.iso8859-1.adapted
etc/rc.d/netatalk.sh.sample
etc/uams/uams_clrtxt.so
etc/uams/uams_dhx.so
%%NETATALKPAM%%etc/uams/uams_dhx_pam.so
%%NETATALKPAM%%etc/uams/uams_dhx_pam.a
etc/uams/uams_dhx_passwd.so
etc/uams/uams_dhx_passwd.a
etc/uams/uams_guest.so
etc/uams/uams_guest.a
%%NETATALKPAM%%etc/uams/uams_pam.so
%%NETATALKPAM%%etc/uams/uams_pam.a
etc/uams/uams_passwd.so
etc/uams/uams_passwd.a
etc/uams/uams_randnum.so
etc/uams/uams_randnum.a
@dirrm etc/nls
@dirrm etc/uams
include/atalk/adouble.h
include/atalk/aep.h
include/atalk/afp.h
include/atalk/asp.h
include/atalk/atp.h
include/atalk/cnid.h
include/atalk/compat.h
include/atalk/ddp.h
include/atalk/dsi.h
include/atalk/nbp.h
include/atalk/netddp.h
include/atalk/pap.h
include/atalk/paths.h
include/atalk/rtmp.h
include/atalk/server_child.h
include/atalk/uam.h
include/atalk/util.h
include/atalk/zip.h
include/netatalk/aarp.c
include/netatalk/aarp.h
include/netatalk/at.h
include/netatalk/at_control.c
include/netatalk/at_proto.c
include/netatalk/at_var.h
include/netatalk/ddp.h
include/netatalk/endian.h
include/netatalk/ddp_var.h
include/netatalk/ddp_input.c
include/netatalk/ddp_output.c
include/netatalk/ddp_usrreq.c
include/netatalk/phase2.h
@dirrm include/atalk
@dirrm include/netatalk
lib/libatalk.a
share/aclocal/netatalk.m4
share/netatalk/pagecount.ps
@dirrm share/netatalk
