# misc/matrix-kmod/pkg-plist
etc/rc.d/matrix_saver.sh
modules/matrix_saver.ko
@unexec rmdir -p %D/modules 2>/dev/null || true
