# games/quakeserver/pkg-plist
etc/rc.d/quakeserver.sh
quakeserver/id1/pak0.pak
quakeserver/unixded
@exec chown quakerun:quakerun %D/quakeserver/unixded
quakeserver/readme.unixded
quakeserver/readme.port
quakeserver/qserver.sh
quakeserver/stopserver.sh
quakeserver/quakeserver.log
@exec chown quakerun:quakerun %D/quakeserver/quakeserver.log
@exec chmod 664 %D/quakeserver/quakeserver.log
@exec mkfifo %D/quakeserver/quakeserver.in
@exec chown quakerun:quakerun %D/quakeserver/quakeserver.in
@exec chmod 660 %D/quakeserver/quakeserver.in
@unexec rm -f %D/quakeserver/quakeserver.in %D/quakeserver/quakeserver.log
@dirrm quakeserver/id1
@dirrm quakeserver
@unexec pw user del quakerun
