# audio/fmio/pkg-plist
bin/fmio
etc/rc.d/bktrctl.sh
share/doc/fmio/Changelog
share/doc/fmio/FAQ
share/doc/fmio/FAQ.ru
share/doc/fmio/README
@dirrm share/doc/fmio
