# audio/fmio/pkg-plist
bin/fmio
etc/rc.d/bktrctl.sh
%%DOCSDIR%%/Changelog
%%DOCSDIR%%/FAQ
%%DOCSDIR%%/FAQ.ru
%%DOCSDIR%%/README
@dirrm %%DOCSDIR%%
