# net/isc-dhcp41-server/pkg-plist
@unexec %D/etc/rc.d/isc-dhcpd forcestop 2>/dev/null || true
%%IPV6%%@unexec %D/etc/rc.d/isc-dhcpd6 forcestop 2>/dev/null || true
@unexec if cmp -s %D/etc/dhcpd.conf.sample %D/etc/dhcpd.conf; then rm -f %D/etc/dhcpd.conf; fi
etc/dhcpd.conf.sample
@exec if [ ! -f %D/etc/dhcpd.conf ] ; then cp -p %D/%F %B/dhcpd.conf; fi
%%IPV6%%@unexec if cmp -s %D/etc/dhcpd6.conf.sample %D/etc/dhcpd6.conf; then rm -f %D/etc/dhcpd6.conf; fi
%%IPV6%%etc/dhcpd6.conf.sample
%%IPV6%%@exec if [ ! -f %D/etc/dhcpd6.conf ] ; then cp -p %D/%F %B/dhcpd6.conf; fi
%%IPV6%%etc/rc.d/isc-dhcpd6
sbin/dhcpd
bin/omshell
%%LDAP%%bin/dhcpd-conf-to-ldap
man/man1/omshell.1.gz
man/man3/dhcpctl.3.gz
man/man3/omapi.3.gz
man/man5/dhcpd.conf.5.gz
man/man5/dhcpd.leases.5.gz
man/man5/dhcp-eval.5.gz
man/man5/dhcp-options.5.gz
man/man8/dhcpd.8.gz
include/dhcpctl.h
include/isc-dhcp/boolean.h
include/isc-dhcp/dst.h
include/isc-dhcp/int.h
include/isc-dhcp/lang.h
include/isc-dhcp/list.h
include/isc-dhcp/result.h
include/isc-dhcp/types.h
include/isc-dhcp/string.h
include/isc-dhcp/print.h
include/isc-dhcp/commandline.h
include/isc-dhcp/mem.h
include/isc-dhcp/formatcheck.h
include/omapip/alloc.h
include/omapip/buffer.h
include/omapip/omapip.h
include/omapip/omapip_p.h
include/omapip/hash.h
include/omapip/trace.h
include/omapip/convert.h
lib/libdhcpctl.a
lib/libomapi.a
lib/libdst.a
%%LDAP%%%%SCHEMA_DIR%%/dhcp.schema
%%LDAP%%@dirrm %%SCHEMA_DIR%%
@dirrm include/omapip
@dirrm include/isc-dhcp
