# irc/anope/files/anope.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: anope
# REQUIRE: LOGIN mysql
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable anope:
# anope_enable=YES
#

. /etc/rc.subr

name=anope
rcvar=anope_enable

load_rc_config $name

: ${anope_enable:=NO}
: ${anope_user:=nobody}

command=%%PREFIX%%/anope/bin/services
pidfile=%%PREFIX%%/anope/data/services.pid
extra_commands=reload

run_rc_command "$1"
