# irc/anope/files/anope.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: anope
# REQUIRE: LOGIN mysql
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable anope:
# anope_enable=YES
#
# Also, to prevent running anope as root You may use
# anope_user="nobody"
# and
# anope_chdir="%%PREFIX%%/share/anope"
#

. /etc/rc.subr

name=anope
rcvar=anope_enable

load_rc_config $name

command=%%PREFIX%%/libexec/anope/services

anope_enable=${anope_enable:-"NO"}

run_rc_command "$1"
