# mail/milter-sender/pkg-plist
@unexec if [ -f %D/etc/milter-sender.cf ] && cmp -s %D/etc/milter-sender.cf %D/etc/milter-sender.cf.sample; then rm -f %D/etc/milter-sender.cf; fi
etc/milter-sender.cf.new
etc/milter-sender.cf.sample
@exec if [ ! -f %D/etc/milter-sender.cf ] ; then cp -p %D/%F %B/milter-sender.cf; fi
@unexec if [ -f %D/etc/rc.d/milter-sender.sh ] && cmp -s %D/etc/rc.d/milter-sender.sh %D/etc/rc.d/milter-sender.sh.sample; then rm -f %D/etc/rc.d/milter-sender.sh; fi
etc/rc.d/milter-sender.sh.sample
@exec if [ ! -f %D/etc/rc.d/milter-sender.sh ] ; then cp %D/%F %B/milter-sender.sh; fi
sbin/milter-sender
@exec install -d -o smmsp -g smmsp -m 755 /var/spool/milter-sender
@unexec rm -rf /var/spool/milter-sender 2>&1 >/dev/null || true
@unexec echo
@unexec echo "If permanently deleting this package, don't forget to delete"
@unexec echo "the inclusion of milter-sender.mc from /etc/mail/YOUR-CONF.mc"
@unexec echo "and then to rebuild and reinstall the sendmail configuration file!"
@unexec echo
@cwd /
@unexec if [ -f %D/etc/mail/milter-sender.mc ] && cmp -s %D/etc/mail/milter-sender.mc %D/etc/mail/milter-sender.mc.sample; then rm -f %D/etc/mail/milter-sender.mc; fi
etc/mail/milter-sender.mc.sample
@exec if [ ! -f %D/etc/mail/milter-sender.mc ] ; then cp %D/%F %B/milter-sender.mc; fi
@exec chgrp smmsp %D/etc/mail/access.db %D/etc/mail/aliases.db
@exec chmod g+r %D/etc/mail/access.db %D/etc/mail/aliases.db
