# mail/milter-sender/pkg-plist
@unexec if [ -f %D/etc/milter-sender.cf ] && cmp -s %D/etc/milter-sender.cf %D/etc/milter-sender.cf.sample; then rm -f %D/etc/milter-sender.cf; fi
etc/milter-sender.cf.sample
@exec if [ ! -f %D/etc/milter-sender.cf ] ; then cp -p %D/%F %B/milter-sender.cf; fi
@unexec if [ -f %D/etc/rc.d/milter-sender.sh ] && cmp -s %D/etc/rc.d/milter-sender.sh %D/etc/rc.d/milter-sender.sh.sample; then rm -f %D/etc/rc.d/milter-sender.sh; fi
etc/rc.d/milter-sender.sh.sample
@exec if [ ! -f %D/etc/rc.d/milter-sender.sh ] ; then cp %D/%F %B/milter-sender.sh; fi
sbin/milter-sender
@exec install -d -o smmsp -g smmsp -m 755 /var/spool/milter-sender
@unexec rm -rf /var/spool/milter-sender
%%PORTDOCS%%%%DOCSDIR%%/CHANGES.TXT
%%PORTDOCS%%%%DOCSDIR%%/Img/bg.gif
%%PORTDOCS%%%%DOCSDIR%%/Img/dilbert2003071742312.gif
%%PORTDOCS%%%%DOCSDIR%%/Img/dilbert2003228690821.gif
%%PORTDOCS%%%%DOCSDIR%%/Img/dilbert20040429.gif
%%PORTDOCS%%%%DOCSDIR%%/Img/mrpixel.gif
%%PORTDOCS%%%%DOCSDIR%%/Img/overview.png
%%PORTDOCS%%%%DOCSDIR%%/Img/sendmail-bat-left-142x38.png
%%PORTDOCS%%%%DOCSDIR%%/Img/sendmail-bat-right-142x38.png
%%PORTDOCS%%%%DOCSDIR%%/Img/sun.gif
%%PORTDOCS%%%%DOCSDIR%%/LICENSE.TXT
%%PORTDOCS%%%%DOCSDIR%%/index.shtml
%%PORTDOCS%%%%DOCSDIR%%/mailto.js
%%PORTDOCS%%%%DOCSDIR%%/style.css
@unexec echo
@unexec echo "If permanently deleting this package, don't forget to delete"
@unexec echo "the inclusion of milter-sender.mc from /etc/mail/YOUR-CONF.mc"
@unexec echo "and then to rebuild and reinstall the sendmail configuration file!"
@unexec echo
%%PORTDOCS%%@dirrm %%DOCSDIR%%/Img
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@cwd /
@unexec if [ -f %D/etc/mail/milter-sender.mc ] && cmp -s %D/etc/mail/milter-sender.mc %D/etc/mail/milter-sender.mc.sample; then rm -f %D/etc/mail/milter-sender.mc; fi
etc/mail/milter-sender.mc.sample
@exec if [ ! -f %D/etc/mail/milter-sender.mc ] ; then cp %D/%F %B/milter-sender.mc; fi
@exec chgrp smmsp %D/etc/mail/access.db %D/etc/mail/aliases.db
@exec chmod g+r %D/etc/mail/access.db %D/etc/mail/aliases.db
