# security/openvas-manager/files/openrc-openvasmd.in
#!/sbin/openrc-run

name=openvasmd
command="%%PREFIX%%/sbin/openvasmd"
pidfile=/var/run/openvasmd.pid

depend() {
		keyword -shutdown
}
