# audio/teamspeak3-server/files/openrc-teamspeak.in
#!/sbin/openrc-run

name="teamspeak"

db_dir=%%DBDIR%%
log_dir=%%LOGDIR%%

pidfile=%%DBDIR%%/teamspeak_server.pid
command=%%PREFIX%%/libexec/ts3server
command_user=teamspeak
supervisor=supervise-daemon
command_args="dbsqlpath=%%PREFIX%%/share/teamspeak/server/sql/ inifile=%%PREFIX%%/etc/teamspeak/ts3server.ini licensepath=%%PREFIX%%/etc/teamspeak/ logpath=$log_dir"
teamspeak_chdir=$db_dir
required_dirs="$db_dir $log_dir"

LD_LIBRARY_PATH=%%PREFIX%%/lib/teamspeak/server:$LD_LIBRARY_PATH
export LD_LIBRARY_PATH

depend() {
	keyword -shutdown
}
