# comms/hylafax/pkg-plist
/var/spool/fax/COPYRIGHT
/var/spool/fax/FIFO
/var/spool/fax/bin/faxrcvd
/var/spool/fax/bin/mkcover
/var/spool/fax/bin/notify
/var/spool/fax/bin/notify.awk
/var/spool/fax/bin/pcl2fax
/var/spool/fax/bin/pollrcvd
/var/spool/fax/bin/ps2fax.dps
/var/spool/fax/bin/ps2fax.gs
/var/spool/fax/bin/ps2fax.imp
/var/spool/fax/bin/tiff2fax
/var/spool/fax/bin/wedged
/var/spool/fax/config/att-dataport
/var/spool/fax/config/att-dataport-2
/var/spool/fax/config/att-dataport-2.0
/var/spool/fax/config/cirrus
/var/spool/fax/config/class1
/var/spool/fax/config/class2
/var/spool/fax/config/class2.0
/var/spool/fax/config/dsi-scout+
/var/spool/fax/config/e-tech
/var/spool/fax/config/everex-2496d
/var/spool/fax/config/exar
/var/spool/fax/config/gvc-288
/var/spool/fax/config/gvc-288-2
/var/spool/fax/config/intel-400e
/var/spool/fax/config/moto-288
/var/spool/fax/config/mt-1432
/var/spool/fax/config/nuvo-voyager
/var/spool/fax/config/ppi-pm14400fxmt
/var/spool/fax/config/ppi-pm14400fxsa
/var/spool/fax/config/rc144ac
/var/spool/fax/config/rc144ac-1
/var/spool/fax/config/rc288dpi
/var/spool/fax/config/rc288dpi-1
/var/spool/fax/config/rc32acl
/var/spool/fax/config/rc32acl-1
/var/spool/fax/config/rockwell
/var/spool/fax/config/skel
/var/spool/fax/config/telebit-qblazer
/var/spool/fax/config/telebit-t3000
/var/spool/fax/config/telebit-wb
/var/spool/fax/config/tricom2842
/var/spool/fax/config/uds-fastalk
/var/spool/fax/config/usr-2.0
/var/spool/fax/config/usr-rts
/var/spool/fax/config/usr-xon
/var/spool/fax/config/zyxel-1496e
/var/spool/fax/config/zyxel-1496e-1
/var/spool/fax/config/zyxel-1496e-2.0
/var/spool/fax/config/zyxel-2864
/var/spool/fax/etc/cover.templ
/var/spool/fax/etc/dialrules
/var/spool/fax/etc/dialrules.europe
/var/spool/fax/etc/dialrules.sf-ba
/var/spool/fax/etc/dpsprinter.ps
/var/spool/fax/etc/hosts
/var/spool/fax/etc/lutRS18.pcf
/var/spool/fax/etc/xferlog
bin/faxalter
bin/faxcover
bin/faxmail
bin/faxrm
bin/faxstat
bin/sendfax
bin/sendpage
etc/rc.d/hylafax.sh.sample
lib/fax/faxcover.ps
lib/fax/faxmail.ps
lib/fax/hfaxd.conf
lib/fax/pagesizes
lib/fax/typerules
man/man1/faxalter.1.gz
man/man1/faxcover.1.gz
man/man1/faxmail.1.gz
man/man1/faxrm.1.gz
man/man1/faxstat.1.gz
man/man1/hylafax.1.gz
man/man1/sendfax.1.gz
man/man1/sendpage.1.gz
man/man1/sgi2fax.1.gz
man/man1/textfmt.1.gz
man/man5/cid.5f.gz
man/man5/config.5f.gz
man/man5/destctrls.5f.gz
man/man5/dialrules.5f.gz
man/man5/hosts.5f.gz
man/man5/hylafax.5f.gz
man/man5/info.5f.gz
man/man5/log.5f.gz
man/man5/pagermap.5f.gz
man/man5/pagesizes.5f.gz
man/man5/recvq.5f.gz
man/man5/sendq.5f.gz
man/man5/shutdown.5f.gz
man/man5/status.5f.gz
man/man5/tsi.5f.gz
man/man5/typerules.5f.gz
man/man5/xferlog.5f.gz
man/man8/choptest.8c.gz
man/man8/cqtest.8c.gz
man/man8/dialtest.8c.gz
man/man8/faxabort.8c.gz
man/man8/faxaddmodem.8c.gz
man/man8/faxanswer.8c.gz
man/man8/faxconfig.8c.gz
man/man8/faxcron.8c.gz
man/man8/faxgetty.8c.gz
man/man8/faxinfo.8c.gz
man/man8/faxmodem.8c.gz
man/man8/faxq.8c.gz
man/man8/faxqclean.8c.gz
man/man8/faxquit.8c.gz
man/man8/faxrcvd.8c.gz
man/man8/faxsend.8c.gz
man/man8/faxsetup.8c.gz
man/man8/faxstate.8c.gz
man/man8/faxwatch.8c.gz
man/man8/hfaxd.8c.gz
man/man8/mkcover.8c.gz
man/man8/notify.8c.gz
man/man8/pagesend.8c.gz
man/man8/pollrcvd.8c.gz
man/man8/ps2fax.8c.gz
man/man8/recvstats.8c.gz
man/man8/tagtest.8c.gz
man/man8/tiff2fax.8c.gz
man/man8/tiffcheck.8c.gz
man/man8/tsitest.8c.gz
man/man8/wedged.8c.gz
man/man8/xferstats.8c.gz
sbin/choptest
sbin/cqtest
sbin/dialtest
sbin/faxabort
sbin/faxaddmodem
sbin/faxanswer
sbin/faxconfig
sbin/faxcron
sbin/faxgetty
sbin/faxinfo
sbin/faxmodem
sbin/faxmsg
sbin/faxq
sbin/faxqclean
sbin/faxquit
sbin/faxsend
sbin/faxsetup
sbin/faxsetup.bsdi
sbin/faxsetup.irix
sbin/faxsetup.linux
sbin/faxstate
sbin/faxwatch
sbin/hfaxd
sbin/lockname
sbin/ondelay
sbin/pagesend
sbin/probemodem
sbin/recvstats
sbin/tagtest
sbin/textfmt
sbin/tiffcheck
sbin/tsitest
sbin/typetest
sbin/xferstats
@dirrm /var/spool/fax
@dirrm /var/spool/fax/archive
@dirrm /var/spool/fax/bin
@dirrm /var/spool/fax/client
@dirrm /var/spool/fax/config
@dirrm /var/spool/fax/dev
@dirrm /var/spool/fax/docq
@dirrm /var/spool/fax/doneq
@dirrm /var/spool/fax/etc
@dirrm /var/spool/fax/info
@dirrm /var/spool/fax/log
@dirrm /var/spool/fax/pollq
@dirrm /var/spool/fax/recvq
@dirrm /var/spool/fax/sendq
@dirrm /var/spool/fax/status
@dirrm /var/spool/fax/tmp
@dirrm lib/fax
