# mail/popa3d/pkg-plist
@comment $FreeBSD$
%%STANDALONE_POP3%%etc/rc.d/popa3d.sh-dist
libexec/popa3d
%%PORTDOCS%%%%DOCSDIR%%/DESIGN
%%PORTDOCS%%%%DOCSDIR%%/LICENSE
%%SMTP_AFTER_POP3%%%%PORTDOCS%%%%DOCSDIR%%/POPAUTH
%%PORTDOCS%%@dirrm %%DOCSDIR%%
