# net/iet/files/ietd.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: ietd
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="ietd"
rcvar=`set_rcvar`

command=%%PREFIX%%/sbin/${name}
required_modules="iet"
stop_precmd="ietd_prestop"

load_rc_config $name

: ${ietd_enable:="NO"}

ietd_prestop()
{
	%%PREFIX%%/sbin/ietadm --op delete
}

run_rc_command "$1"
