# www/jakarta-tomcat3/pkg-plist
@comment $FreeBSD$
@unexec if cmp -s %D/%%T%%/conf/jk/mod_jk.conf %D/%%T%%/conf/jk/mod_jk.conf.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/jk/mod_jk.conf; fi
@unexec if cmp -s %D/%%T%%/conf/jk/mod_jk.conf-jni %D/%%T%%/conf/jk/mod_jk.conf-jni.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/jk/mod_jk.conf-jni; fi
@unexec if cmp -s %D/%%T%%/conf/jk/obj.conf %D/%%T%%/conf/jk/obj.conf.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/jk/obj.conf; fi
@unexec if cmp -s %D/%%T%%/conf/jk/uriworkermap.properties %D/%%T%%/conf/jk/uriworkermap.properties.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/jk/uriworkermap.properties; else fi
@unexec if cmp -s %D/%%T%%/conf/jk/workers.properties %D/%%T%%/conf/jk/workers.properties.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/jk/workers.properties; fi
@unexec if cmp -s %D/%%T%%/conf/jk/wrapper.properties %D/%%T%%/conf/jk/wrapper.properties.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/jk/wrapper.properties; fi
@unexec if cmp -s %D/%%T%%/conf/jserv/tomcat.conf %D/%%T%%/conf/jserv/tomcat.conf.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/jserv/tomcat.conf; fi
@unexec if cmp -s %D/%%T%%/conf/jserv/tomcat.properties %D/%%T%%/conf/jserv/tomcat.properties.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/jserv/tomcat.properties; fi
@unexec if cmp -s %D/%%T%%/conf/jk2.properties %D/%%T%%/conf/jk2.properties.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/jk2.properties; fi
@unexec if cmp -s %D/%%T%%/conf/apps-127.0.0.1.xml %D/%%T%%/conf/apps-127.0.0.1.xml.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/apps-127.0.0.1.xml; fi
@unexec if cmp -s %D/%%T%%/conf/apps-admin.xml %D/%%T%%/conf/apps-admin.xml.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/apps-admin.xml; fi
@unexec if cmp -s %D/%%T%%/conf/apps-examples.xml %D/%%T%%/conf/apps-examples.xml.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/apps-examples.xml; fi
@unexec if cmp -s %D/%%T%%/conf/tomcat.policy %D/%%T%%/conf/tomcat.policy.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/tomcat.policy; fi
@unexec if cmp -s %D/%%T%%/conf/modules.xml %D/%%T%%/conf/modules.xml.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/modules.xml; fi
@unexec if cmp -s %D/%%T%%/conf/server.xml %D/%%T%%/conf/server.xml.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/server.xml; fi
@unexec if cmp -s %D/%%T%%/conf/users/admin-users.xml %D/%%T%%/conf/users/admin-users.xml.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/users/admin-users.xml; fi
@unexec if cmp -s %D/%%T%%/conf/users/example-users.xml %D/%%T%%/conf/users/example-users.xml.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/users/example-users.xml; fi
@unexec if cmp -s %D/%%T%%/conf/users/global-users.xml %D/%%T%%/conf/users/global-users.xml.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/users/global-users.xml; fi
@unexec if cmp -s %D/%%T%%/conf/users/tomcat-users.xml %D/%%T%%/conf/users/tomcat-users.xml.%%CONF_EXT%%; then rm -f %D/%%T%%/conf/users/tomcat-users.xml; fi
bin/tomcat3ctl
etc/rc.d/020.jakarta-tomcat3.sh
%%T%%/logs/stdout.log
%%T%%/logs/stderr.log
%%T%%/doc/images/banner.gif
%%T%%/doc/images/tomcat.gif
%%T%%/doc/appdev/sample/etc/web.xml
%%T%%/doc/appdev/sample/lib/README
%%T%%/doc/appdev/sample/src/Hello.java
%%T%%/doc/appdev/sample/web/images/tomcat.gif
%%T%%/doc/appdev/sample/web/hello.jsp
%%T%%/doc/appdev/sample/web/index.html
%%T%%/doc/appdev/sample/build.sh
%%T%%/doc/appdev/sample/build.xml
%%T%%/doc/appdev/build.xml.txt
%%T%%/doc/appdev/contents.html
%%T%%/doc/appdev/deployment.html
%%T%%/doc/appdev/footer.html
%%T%%/doc/appdev/header.html
%%T%%/doc/appdev/index.html
%%T%%/doc/appdev/installation.html
%%T%%/doc/appdev/introduction.html
%%T%%/doc/appdev/processes.html
%%T%%/doc/appdev/source.html
%%T%%/doc/appdev/tomcat.gif
%%T%%/doc/appdev/web.xml.txt
%%T%%/doc/JDBCRealm-howto.html
%%T%%/doc/NT-Service-howto.html
%%T%%/doc/Tomcat-Workers-HowTo.html
%%T%%/doc/Tomcat-on-NetWare-HowTo.html
%%T%%/doc/faq
%%T%%/doc/in-process-howto.html
%%T%%/doc/index.html
%%T%%/doc/internal.html
%%T%%/doc/mod_jk-howto.html
%%T%%/doc/perfNumbers.txt
%%T%%/doc/readme
%%T%%/doc/readme-3.3.1
%%T%%/doc/readme-3.3.2
%%T%%/doc/serverxml.html
%%T%%/doc/style.css
%%T%%/doc/tomcat-apache-howto.html
%%T%%/doc/tomcat-iis-howto.html
%%T%%/doc/tomcat-netscape-howto.html
%%T%%/doc/tomcat-security.html
%%T%%/doc/tomcat-ssl-howto.html
%%T%%/doc/tomcat-ug.html
%%T%%/doc/AJPv13.html
%%T%%/native/mod_jserv/mod_jserv.c
%%T%%/native/mod_jserv/Makefile.libdir
%%T%%/native/mod_jserv/Makefile.tmpl
%%T%%/native/mod_jserv/Makefile.win32
%%T%%/native/mod_jserv/autochange.c
%%T%%/native/mod_jserv/jserv.h
%%T%%/native/mod_jserv/jserv_ajpv11.c
%%T%%/native/mod_jserv/jserv_ajpv12.c
%%T%%/native/mod_jserv/jserv_balance.c
%%T%%/native/mod_jserv/jserv_image.c
%%T%%/native/mod_jserv/jserv_mmap.c
%%T%%/native/mod_jserv/jserv_protocols.c
%%T%%/native/mod_jserv/jserv_status.c
%%T%%/native/mod_jserv/jserv_utils.c
%%T%%/native/mod_jserv/jserv_watchdog.c
%%T%%/native/mod_jserv/jserv_wrapper.c
%%T%%/native/mod_jserv/jserv_wrapper.h
%%T%%/native/mod_jserv/jserv_wrapper_unix.c
%%T%%/native/mod_jserv/jserv_wrapper_win.c
%%T%%/native/mod_jserv/libjserv.module
%%T%%/native/mod_jserv/Makefile.am
%%T%%/lib/apps/README
%%T%%/lib/common/tomcat_core.jar
%%T%%/lib/common/commons-logging-api.jar
%%T%%/lib/common/core_util.jar
%%T%%/lib/common/etomcat.jar
%%T%%/lib/common/jasper-runtime.jar
%%T%%/lib/common/servlet.jar
%%T%%/lib/common/tomcat-util.jar
%%T%%/lib/common/tomcat33-resource.jar
%%T%%/lib/common/README
%%T%%/lib/container/xalan.jar
%%T%%/lib/container/container_util.jar
%%T%%/lib/container/crimson.jar
%%T%%/lib/container/facade22.jar
%%T%%/lib/container/jakarta-regexp-1.3.jar
%%T%%/lib/container/jasper.jar
%%T%%/lib/container/mx4j-jmx.jar
%%T%%/lib/container/mx4j-tools.jar
%%T%%/lib/container/tomcat-coyote.jar
%%T%%/lib/container/tomcat-http11.jar
%%T%%/lib/container/tomcat-jk2.jar
%%T%%/lib/container/tomcat33-coyote.jar
%%T%%/lib/container/tomcat_modules.jar
%%T%%/lib/container/README
%%T%%/lib/endorsed/README
%%T%%/lib/stop-tomcat.jar
%%T%%/lib/tomcat.jar
%%T%%/webapps/ROOT.war
%%T%%/webapps/admin.war
%%T%%/webapps/examples.war
%%T%%/conf/jk/mod_jk.conf.%%CONF_EXT%%
%%T%%/conf/jk/mod_jk.conf-jni.%%CONF_EXT%%
%%T%%/conf/jk/obj.conf.%%CONF_EXT%%
%%T%%/conf/jk/uriworkermap.properties.%%CONF_EXT%%
%%T%%/conf/jk/workers.properties.%%CONF_EXT%%
%%T%%/conf/jk/wrapper.properties.%%CONF_EXT%%
%%T%%/conf/jserv/tomcat.conf.%%CONF_EXT%%
%%T%%/conf/jserv/tomcat.properties.%%CONF_EXT%%
%%T%%/conf/users/admin-users.xml.%%CONF_EXT%%
%%T%%/conf/users/example-users.xml.%%CONF_EXT%%
%%T%%/conf/users/global-users.xml.%%CONF_EXT%%
%%T%%/conf/users/tomcat-users.xml.%%CONF_EXT%%
%%T%%/conf/jk2.properties.%%CONF_EXT%%
%%T%%/conf/modules.xml.%%CONF_EXT%%
%%T%%/conf/server.xml.%%CONF_EXT%%
%%T%%/conf/tomcat.policy.%%CONF_EXT%%
%%T%%/conf/apps-127.0.0.1.xml.%%CONF_EXT%%
%%T%%/conf/apps-admin.xml.%%CONF_EXT%%
%%T%%/conf/apps-examples.xml.%%CONF_EXT%%
%%T%%/bin/jspc
%%T%%/bin/shutdown
%%T%%/bin/startup
%%T%%/bin/tomcat
%%T%%/ant/tomcat-ant.jar
%%T%%/LICENSE
%%T%%/KEYS
@exec mkdir -p %D/%%T%%/conf/auto || true
@exec mkdir -p %D/%%T%%/modules || true
@exec chown -R %%WWWOWN%%:%%WWWGRP%% %D/%%T%%/
@unexec rmdir %D/%%T%%/conf/users 2>/dev/null || true
@unexec rmdir %D/%%T%%/conf/auto 2>/dev/null || true
@unexec rmdir %D/%%T%%/conf/jserv 2>/dev/null || true
@unexec rmdir %D/%%T%%/conf/jk 2>/dev/null || true
@unexec rmdir %D/%%T%%/conf 2>/dev/null || true
@dirrm %%T%%/logs
@dirrm %%T%%/doc/images
@dirrm %%T%%/doc/appdev/sample/etc
@dirrm %%T%%/doc/appdev/sample/lib
@dirrm %%T%%/doc/appdev/sample/src
@dirrm %%T%%/doc/appdev/sample/web/images
@dirrm %%T%%/doc/appdev/sample/web
@dirrm %%T%%/doc/appdev/sample
@dirrm %%T%%/doc/appdev
@dirrm %%T%%/doc
@dirrm %%T%%/native/mod_jserv
@dirrm %%T%%/native
@dirrm %%T%%/modules
@dirrm %%T%%/lib/apps
@dirrm %%T%%/lib/common
@dirrm %%T%%/lib/container
@dirrm %%T%%/lib/endorsed
@dirrm %%T%%/lib
@dirrm %%T%%/webapps
@dirrm %%T%%/bin
@dirrm %%T%%/ant
@unexec rmdir %D/%%T%% 2>/dev/null || true
