# net/splatd/files/openrc-splatd.in
#!/sbin/openrc-run

name=splatd

command=%%PREFIX%%/sbin/splatd
command_args="-f %%PREFIX%%/etc/splat.conf -p /var/run/$name.pid"
pidfile="/var/run/$name.pid"
required_files=%%PREFIX%%/etc/splat.conf

depend() {
	keyword -shutdown
}
