# net/splatd/files/splatd.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: splatd 
# REQUIRE: DAEMON
# BEFORE:  LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name=splatd
rcvar=splatd_enable

command=%%PREFIX%%/sbin/splatd
command_interpreter="%%PYTHON_CMD%%"
command_args="-f %%PREFIX%%/etc/splat.conf -p /var/run/$name.pid"
pidfile="/var/run/$name.pid"
required_files=%%PREFIX%%/etc/splat.conf

# set default
splatd_enable=${splatd_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"
