# databases/firebird-devel/pkg-plist
@comment $FreeBSD$
@unexec rm -f %%LOCALBASE%%/etc/rc.d/000.firebird.sh
firebird/README
firebird/UDF/fbudf.so
firebird/UDF/fbudf.sql
firebird/UDF/ib_udf.so
firebird/UDF/ib_udf.sql
firebird/WhatsNew
@unexec if cmp -s %D/firebird/aliases.conf %D/firebird/aliases.conf.sample; then rm -f %D/firebird/aliases.conf; fi
firebird/aliases.conf.sample
@exec [ -f %B/aliases.conf ] || cp %B/%f %B/aliases.conf
firebird/bin/fb_inet_server
firebird/bin/fb_lock_mgr
firebird/bin/fb_lock_print
firebird/bin/gbak
firebird/bin/gdef
firebird/bin/gds_drop
firebird/bin/gfix
firebird/bin/gpre
firebird/bin/gsec
firebird/bin/gsplit
firebird/bin/gstat
firebird/bin/isql
firebird/bin/qli
%%PORTDOCS%%%%DOCSDIR%%/README.NTSecurity
%%PORTDOCS%%%%DOCSDIR%%/README.Optimizer.txt
%%PORTDOCS%%%%DOCSDIR%%/README.Win32LibraryInstallation.txt
%%PORTDOCS%%%%DOCSDIR%%/README.build.mingw.html
%%PORTDOCS%%%%DOCSDIR%%/README.build.msvc.html
%%PORTDOCS%%%%DOCSDIR%%/README.coding.style
%%PORTDOCS%%%%DOCSDIR%%/README.intl
%%PORTDOCS%%%%DOCSDIR%%/README.instsvc
%%PORTDOCS%%%%DOCSDIR%%/README.isql_enhancements.txt
%%PORTDOCS%%%%DOCSDIR%%/README.makefiles
%%PORTDOCS%%%%DOCSDIR%%/README.performance_monitoring
%%PORTDOCS%%%%DOCSDIR%%/README.sha1
%%PORTDOCS%%%%DOCSDIR%%/README.user
%%PORTDOCS%%%%DOCSDIR%%/README.user.embedded
%%PORTDOCS%%%%DOCSDIR%%/README.user.troubleshooting
%%PORTDOCS%%%%DOCSDIR%%/README.xnet
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.PSQL_stack_trace.txt
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.aggregate_tracking
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.case
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.coalesce
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.context_variables
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.context_variables2
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.cursors
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.data_type_results_of_aggregations.txt
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.data_types
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.ddl.txt
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.default_parameters
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.derived_tables.txt
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.distinct
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.exception_handling
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.execute_block
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.execute_statement
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.explicit_locks
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.expression_indices
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.iif
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.keywords
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.leave_labels
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.length
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.null_value
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.nullif
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.order_by_expressions_nulls
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.plan
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.returning
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.rows
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.savepoints
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.select_expressions
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.sequence_generators
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.trim
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.universal_triggers
%%PORTDOCS%%%%DOCSDIR%%/sql.extensions/README.view_updates
%%PORTDOCS%%@dirrm %%DOCSDIR%%/sql.extensions
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@unexec if cmp -s %D/firebird/firebird.conf %D/firebird/firebird.conf.sample; then rm -f %D/firebird/firebird.conf; fi
firebird/firebird.conf.sample
@exec [ -f %B/firebird.conf ] || cp %B/%f %B/firebird.conf
@unexec rm -f %D/firebird/firebird.log
firebird/firebird.msg
firebird/help/help.fdb
firebird/include/blr.h
firebird/include/ib_util.h
firebird/include/ibase.h
firebird/include/iberror.h
firebird/include/perf.h
firebird/intl/fbintl
@unexec rm -f %D/firebird/isc_event1.`hostname`
@unexec rm -f %D/firebird/isc_init1.`hostname`
@unexec rm -f %D/firebird/isc_lock1.`hostname`
firebird/lib/libfbclient.so
firebird/lib/libfbclient.so.2
firebird/lib/libfbclient.so.2.0.0
firebird/lib/libfbembed.so
firebird/lib/libfbembed.so.2
firebird/lib/libfbembed.so.2.0.0
firebird/lib/libib_util.so
firebird/security2.fdb.sample
firebird/security2.fdb
@dirrm firebird/UDF
@dirrm firebird/bin
@dirrm firebird/help
@dirrm firebird/include
@dirrm firebird/intl
@dirrm firebird/lib
@unexec rmdir %D/firebird 2>/dev/null || true
