# databases/firebird-devel/pkg-plist
@comment $FreeBSD$
@unexec rm -f %%LOCALBASE%%/etc/rc.d/000.firebird.sh
README
UDF/fbudf.so
UDF/fbudf.sql
UDF/ib_udf.so
UDF/ib_udf.sql
WhatsNew
aliases.conf
aliases.conf.sample
bin/fb_inet_server
bin/fb_lock_mgr
bin/fb_lock_print
bin/gbak
bin/gdef
bin/gds_drop
bin/gfix
bin/gpre
bin/gsec
bin/gsplit
bin/gstat
bin/isql
bin/qli
doc/README.NTSecurity
doc/README.Optimizer.txt
doc/README.Win32LibraryInstallation.txt
doc/README.build.mingw.html
doc/README.build.msvc.html
doc/README.coding.style
doc/README.instsvc
doc/README.isql_enhancements.txt
doc/README.makefiles
doc/README.performance_monitoring
doc/README.sha1
doc/README.user
doc/README.user.embedded
doc/README.user.troubleshooting
doc/README.xnet
doc/sql.extensions/README.PSQL_stack_trace.txt
doc/sql.extensions/README.aggregate_tracking
doc/sql.extensions/README.case
doc/sql.extensions/README.coalesce
doc/sql.extensions/README.context_variables
doc/sql.extensions/README.context_variables2
doc/sql.extensions/README.cursors
doc/sql.extensions/README.data_type_results_of_aggregations.txt
doc/sql.extensions/README.data_types
doc/sql.extensions/README.ddl.txt
doc/sql.extensions/README.default_parameters
doc/sql.extensions/README.derived_tables.txt
doc/sql.extensions/README.distinct
doc/sql.extensions/README.exception_handling
doc/sql.extensions/README.execute_block
doc/sql.extensions/README.execute_statement
doc/sql.extensions/README.explicit_locks
doc/sql.extensions/README.expression_indices
doc/sql.extensions/README.iif
doc/sql.extensions/README.keywords
doc/sql.extensions/README.leave_labels
doc/sql.extensions/README.null_value
doc/sql.extensions/README.nullif
doc/sql.extensions/README.order_by_expressions_nulls
doc/sql.extensions/README.plan
doc/sql.extensions/README.rows
doc/sql.extensions/README.savepoints
doc/sql.extensions/README.select_expressions
doc/sql.extensions/README.sequence_generators
doc/sql.extensions/README.universal_triggers
doc/sql.extensions/README.view_updates
firebird.conf
firebird.conf.sample
@unexec rm -f %D/firebird.log
firebird.msg
help/help.fdb
include/blr.h
include/ib_util.h
include/ibase.h
include/iberror.h
include/perf.h
intl/fbintl
@unexec rm -f %D/isc_event1.`hostname`
@unexec rm -f %D/isc_init1.`hostname`
@unexec rm -f %D/isc_lock1.`hostname`
lib/libfbclient.so
lib/libfbclient.so.2
lib/libfbclient.so.2.0.0
lib/libfbembed.so
lib/libfbembed.so.2
lib/libfbembed.so.2.0.0
lib/libib_util.so
security2.fdb.sample
security2.fdb
@dirrm UDF
@dirrm doc/sql.extensions
@dirrm doc
@dirrm help
@dirrm intl
