# databases/firebird-devel/pkg-plist
etc/rc.d/000.firebird.sh
firebird/README
firebird/UDF/fbudf.so
firebird/UDF/fbudf.sql
firebird/UDF/ib_udf.so
firebird/UDF/ib_udf.sql
firebird/WhatsNew
firebird/aliases.conf
firebird/bin/fb_inet_server
firebird/bin/fb_lock_mgr
firebird/bin/fb_lock_print
firebird/bin/gbak
firebird/bin/gdef
firebird/bin/gds_drop
firebird/bin/gfix
firebird/bin/gpre
firebird/bin/gsec
firebird/bin/gsplit
firebird/bin/gstat
firebird/bin/isql
firebird/bin/qli
firebird/doc/README.NTSecurity
firebird/doc/README.build.mingw.html
firebird/doc/README.build.msvc.html
firebird/doc/README.coding.style
firebird/doc/README.instsvc
firebird/doc/README.makefiles
firebird/doc/README.user
firebird/doc/README.user.embedded
firebird/doc/README.user.troubleshooting
firebird/doc/sql.extensions/README.aggregate_tracking
firebird/doc/sql.extensions/README.case
firebird/doc/sql.extensions/README.coalesce
firebird/doc/sql.extensions/README.context_variables
firebird/doc/sql.extensions/README.data_types
firebird/doc/sql.extensions/README.exception_handling
firebird/doc/sql.extensions/README.execute_statement
firebird/doc/sql.extensions/README.explicit_locks
firebird/doc/sql.extensions/README.nullif
firebird/doc/sql.extensions/README.order_by_expressions_nulls
firebird/doc/sql.extensions/README.savepoints
firebird/doc/sql.extensions/README.universal_triggers
firebird/firebird.conf
@unexec rm -f %D/firebird/firebird.log
firebird/firebird.msg
firebird/help/help.fdb
firebird/include/blr.h
firebird/include/gds.h
firebird/include/ib_util.h
firebird/include/ibase.h
firebird/include/iberror.h
firebird/include/perf.h
firebird/intl/fbintl
@unexec rm -f %D/firebird/isc_event1.`hostname`
@unexec rm -f %D/firebird/isc_init1.`hostname`
@unexec rm -f %D/firebird/isc_lock1.`hostname`
firebird/lib/libfbclient.so
firebird/lib/libfbclient.so.2
firebird/lib/libfbclient.so.2.0.0
firebird/lib/libfbembed.so
firebird/lib/libfbembed.so.2
firebird/lib/libfbembed.so.2.0.0
firebird/lib/libib_util.so
firebird/security.fdb
@dirrm firebird/UDF
@dirrm firebird/bin
@dirrm firebird/doc/sql.extensions
@dirrm firebird/doc
@dirrm firebird/help
@dirrm firebird/include
@dirrm firebird/intl
@dirrm firebird/lib
@dirrm firebird
