# sysutils/runit/files/runsvdir.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: /service
# BEFORE:  SERVERS
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable `runsvdir':
#
#runsvdir_enable="YES"		(default: NO)
#runsvdir_path="/service"	(default: %%SERVICE_DIR%%)
#

. /etc/rc.subr

name="runsvdir"
rcvar=runsvdir_enable

command="%%PREFIX%%/sbin/runsvdir"

start_cmd=runsvdir_start
stop_postcmd=runsvdir_post_stop

runsvdir_start()
{
	/bin/test -d ${runsvdir_path} || /bin/mkdir -p ${runsvdir_path}
	/usr/bin/env -i \
	PATH=/command:%%PREFIX%%/sbin:%%PREFIX%%/bin:/sbin:/bin:/usr/sbin:/usr/bin \
	/usr/sbin/daemon \
	${command} ${runsvdir_path}
}

runsvdir_post_stop()
{
	%%PREFIX%%/sbin/sv exit ${runsvdir_path}/* 
}

# read configuration and set defaults
load_rc_config "$name"
: ${runsvdir_enable="NO"}
: ${runsvdir_path="%%SERVICE_DIR%%"}

run_rc_command "$1"
