# security/razorback-masterNugget/files/masterNugget_safed.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: masterNugget_safed
# REQUIRE: masterNugget
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable the razorback master nugget safed:
#
# masterNugget_safed_enable="YES"

. /etc/rc.subr

name=masterNugget_safed
rcvar=masterNugget_safed_enable

load_rc_config $name

masterNugget_safed_enable=${masterNugget_safed_enable:-"NO"}

command="%%PREFIX%%/bin/masterNugget_safed"

pidfile="/var/run/masterNugget_safed.pid"
start_cmd=masterNugget_safed_start
procname="/bin/sh"

masterNugget_safed_start()
{
        /usr/sbin/daemon -cf -p ${pidfile} ${command}
}

run_rc_command "$1"
