# security/razorback-masterNugget/files/openrc-masterNugget_safed.in
#!/sbin/openrc-run

name=masterNugget_safed

command="%%PREFIX%%/bin/masterNugget_safed"
supervisor=supervise-daemon
pidfile="/var/run/masterNugget_safed.pid"

depend() {
	need masterNuggetd
	keyword -shutdown
}
