# audio/mpio-driver/pkg-plist
etc/rc.d/mpio-driver.sh.sample
@unexec rm -f %%MODDEST%%/modules/mpio.ko 2> /dev/null || true
