# net-mgmt/smokeping/pkg-plist
bin/smokeping
bin/tSmoke
etc/rc.d/smokeping.sh
@unexec if cmp -s %D/etc/smokeping/basepage.html %D/etc/smokeping/basepage.html.dist; then rm -f %D/etc/smokeping/basepage.html; fi
etc/smokeping/basepage.html.dist
@exec [ -f %B/basepage.html ] || cp %B/%f %B/basepage.html
@unexec if cmp -s %D/etc/smokeping/config %D/etc/smokeping/config.dist; then rm -f %D/etc/smokeping/config; fi
etc/smokeping/config.dist
@exec [ -f %B/config ] || cp %B/%f %B/config
@unexec if cmp -s %D/etc/smokeping/smokemail %D/etc/smokeping/smokemail.dist; then rm -f %D/etc/smokeping/smokemail; fi
etc/smokeping/smokemail.dist
@exec [ -f %B/smokemail ] || cp %B/%f %B/smokemail
@unexec if cmp -s %D/etc/smokeping/tmail %D/etc/smokeping/tmail.dist; then rm -f %D/etc/smokeping/tmail; fi
etc/smokeping/tmail.dist
@exec [ -f %B/tmail ] || cp %B/%f %B/tmail
%%PORTDOCS%%%%DOCSDIR%%/CHANGES
%%PORTDOCS%%%%DOCSDIR%%/CONTRIBUTORS
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%%%DOCSDIR%%/Smokeping.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping.cgi.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping.cgi.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping_config.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping_config.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping_install.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping_install.txt
%%PORTDOCS%%%%DOCSDIR%%/Config/Grammar.html
%%PORTDOCS%%%%DOCSDIR%%/Config/Grammar.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/matchers/Avgratio.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/matchers/Avgratio.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/matchers/base.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/matchers/base.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/matchers/Median.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/matchers/Median.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/AnotherDNS.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/AnotherDNS.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/AnotherSSH.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/AnotherSSH.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/CiscoRTTMonDNS.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/CiscoRTTMonDNS.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/CiscoRTTMonEchoICMP.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/CiscoRTTMonEchoICMP.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/CiscoRTTMonTcpConnect.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/CiscoRTTMonTcpConnect.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/Curl.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/Curl.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/DNS.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/DNS.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPing.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPing.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPingChargen.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPingChargen.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPingDiscard.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPingDiscard.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPingHttp.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPingHttp.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPingHttps.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPingHttps.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPingIcp.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPingIcp.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPingSmtp.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/EchoPingSmtp.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/FPing.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/FPing.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/FPing6.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/FPing6.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/IOSPing.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/IOSPing.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/LDAP.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/LDAP.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/Radius.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/Radius.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/RemoteFPing.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/RemoteFPing.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/SSH.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/SSH.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/TelnetIOSPing.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/TelnetIOSPing.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/base.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/base.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/basefork.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/basefork.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/basevars.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/basevars.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/passwordchecker.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/passwordchecker.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/skel.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/probes/skel.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/Examples.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/Examples.txt
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/RRDtools.html
%%PORTDOCS%%%%DOCSDIR%%/Smokeping/RRDtools.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping_examples.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping_examples.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping_extend.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping_extend.txt
%%PORTDOCS%%%%DOCSDIR%%/smokeping_upgrade.html
%%PORTDOCS%%%%DOCSDIR%%/smokeping_upgrade.txt
%%PORTDOCS%%%%DOCSDIR%%/tSmoke.html
%%PORTDOCS%%%%DOCSDIR%%/tSmoke.txt
%%PORTDOCS%%%%EXAMPLESDIR%%/config.echoping
%%PORTDOCS%%%%EXAMPLESDIR%%/config.fping-instances
%%PORTDOCS%%%%EXAMPLESDIR%%/config.multiple-probes
%%PORTDOCS%%%%EXAMPLESDIR%%/config.simple
%%PORTDOCS%%%%EXAMPLESDIR%%/config.targetvars-with-Curl
%%PORTDOCS%%%%EXAMPLESDIR%%/config.template
smokeping/htdocs/smokeping.cgi
smokeping/lib/Smokeping.pm
smokeping/lib/Smokeping/ciscoRttMonMIB.pm
smokeping/lib/Smokeping/matchers/base.pm
smokeping/lib/Smokeping/matchers/CheckLatency.pm
smokeping/lib/Smokeping/matchers/CheckLoss.pm
smokeping/lib/Smokeping/matchers/Median.pm
smokeping/lib/Smokeping/matchers/Avgratio.pm
smokeping/lib/Smokeping/probes/AnotherDNS.pm
smokeping/lib/Smokeping/probes/AnotherSSH.pm
smokeping/lib/Smokeping/probes/CiscoRTTMonDNS.pm
smokeping/lib/Smokeping/probes/CiscoRTTMonEchoICMP.pm
smokeping/lib/Smokeping/probes/CiscoRTTMonTcpConnect.pm
smokeping/lib/Smokeping/probes/Curl.pm
smokeping/lib/Smokeping/probes/DNS.pm
smokeping/lib/Smokeping/probes/EchoPing.pm
smokeping/lib/Smokeping/probes/EchoPingChargen.pm
smokeping/lib/Smokeping/probes/EchoPingDiscard.pm
smokeping/lib/Smokeping/probes/EchoPingHttp.pm
smokeping/lib/Smokeping/probes/EchoPingHttps.pm
smokeping/lib/Smokeping/probes/EchoPingIcp.pm
smokeping/lib/Smokeping/probes/EchoPingSmtp.pm
smokeping/lib/Smokeping/probes/FPing.pm
smokeping/lib/Smokeping/probes/FPing6.pm
smokeping/lib/Smokeping/probes/IOSPing.pm
smokeping/lib/Smokeping/probes/LDAP.pm
smokeping/lib/Smokeping/probes/Radius.pm
smokeping/lib/Smokeping/probes/RemoteFPing.pm
smokeping/lib/Smokeping/probes/SSH.pm
smokeping/lib/Smokeping/probes/TacacsPlus.pm
smokeping/lib/Smokeping/probes/TelnetIOSPing.pm
smokeping/lib/Smokeping/probes/base.pm
smokeping/lib/Smokeping/probes/basefork.pm
smokeping/lib/Smokeping/probes/basevars.pm
smokeping/lib/Smokeping/probes/passwordchecker.pm
smokeping/lib/Smokeping/probes/skel.pm
smokeping/lib/Smokeping/Examples.pm
smokeping/lib/Smokeping/RRDtools.pm
smokeping/lib/Config/Grammar.pm
@dirrmtry etc/smokeping
@exec mkdir -p %D/smokeping/htdocs/img; chown www:www %D/smokeping/htdocs/img
@unexec rm -rf %D/smokeping/htdocs/img 2>/dev/null || true
@exec mkdir -p %D/var/smokeping
@exec chown smokeping:smokeping %D/var/smokeping
@dirrmtry var/smokeping
@dirrmtry var
@dirrm smokeping/htdocs
@dirrm smokeping/lib/Config
@dirrm smokeping/lib/Smokeping/probes
@dirrm smokeping/lib/Smokeping/matchers
@dirrm smokeping/lib/Smokeping
@dirrm smokeping/lib
%%PORTDOCS%%@dirrm %%DOCSDIR%%/Config
%%PORTDOCS%%@dirrm %%DOCSDIR%%/Smokeping/probes
%%PORTDOCS%%@dirrm %%DOCSDIR%%/Smokeping/matchers
%%PORTDOCS%%@dirrm %%DOCSDIR%%/Smokeping
%%PORTDOCS%%@dirrm %%DOCSDIR%%
%%PORTDOCS%%@dirrm %%EXAMPLESDIR%%
@dirrmtry smokeping
