# net-mgmt/smokeping/pkg-plist
bin/smokeping
etc/rc.d/smokeping.sh
@unexec if cmp -s %D/etc/smokeping/basepage.html %D/etc/smokeping/basepage.html.dist; then rm -f %D/etc/smokeping/basepage.html; fi
etc/smokeping/basepage.html.dist
@exec [ -f %B/basepage.html ] || cp %B/%f %B/basepage.html
@unexec if cmp -s %D/etc/smokeping/config %D/etc/smokeping/config.dist; then rm -f %D/etc/smokeping/config; fi
etc/smokeping/config.dist
@exec [ -f %B/config ] || cp %B/%f %B/config
@unexec if cmp -s %D/etc/smokeping/smokemail %D/etc/smokeping/smokemail.dist; then rm -f %D/etc/smokeping/smokemail; fi
etc/smokeping/smokemail.dist
@exec [ -f %B/smokemail ] || cp %B/%f %B/smokemail
@unexec if cmp -s %D/etc/smokeping/config-echoping %D/etc/smokeping/config-echoping.dist; then rm -f %D/etc/smokeping/config-echoping; fi
etc/smokeping/config-echoping.dist
@exec [ -f %B/config-echoping ] || cp %B/%f %B/config-echoping
share/doc/smokeping/CHANGES
share/doc/smokeping/CONTRIBUTORS
share/doc/smokeping/COPYING
share/doc/smokeping/COPYRIGHT
share/doc/smokeping/README
share/doc/smokeping/TODO
share/doc/smokeping/EchoPing.pm.html
share/doc/smokeping/EchoPing.pm.txt
share/doc/smokeping/EchoPingChargen.pm.html
share/doc/smokeping/EchoPingChargen.pm.txt
share/doc/smokeping/EchoPingDiscard.pm.html
share/doc/smokeping/EchoPingDiscard.pm.txt
share/doc/smokeping/EchoPingHttp.pm.html
share/doc/smokeping/EchoPingHttp.pm.txt
share/doc/smokeping/EchoPingHttps.pm.html
share/doc/smokeping/EchoPingHttps.pm.txt
share/doc/smokeping/EchoPingIcp.pm.html
share/doc/smokeping/EchoPingIcp.pm.txt
share/doc/smokeping/EchoPingSmtp.pm.html
share/doc/smokeping/EchoPingSmtp.pm.txt
share/doc/smokeping/FPing.pm.html
share/doc/smokeping/FPing.pm.txt
share/doc/smokeping/FPing6.pm.html
share/doc/smokeping/FPing6.pm.txt
share/doc/smokeping/IOSPing.pm.html
share/doc/smokeping/IOSPing.pm.txt
share/doc/smokeping/ParseConfig.pm.html
share/doc/smokeping/ParseConfig.pm.txt
share/doc/smokeping/Smokeping.pm.html
share/doc/smokeping/Smokeping.pm.txt
share/doc/smokeping/base.pm.html
share/doc/smokeping/base.pm.txt
share/doc/smokeping/basefork.pm.html
share/doc/smokeping/basefork.pm.txt
share/doc/smokeping/basevars.pm.html
share/doc/smokeping/basevars.pm.txt
share/doc/smokeping/smokeping.cgi.html
share/doc/smokeping/smokeping.cgi.txt
share/doc/smokeping/smokeping.html
share/doc/smokeping/smokeping.txt
share/doc/smokeping/smokeping_config.html
share/doc/smokeping/smokeping_config.txt
share/doc/smokeping/smokeping_install.html
share/doc/smokeping/smokeping_install.txt
smokeping/htdocs/smokeping.cgi
smokeping/lib/Smokeping.pm
smokeping/lib/probes/FPing.pm
smokeping/lib/probes/FPing6.pm
smokeping/lib/probes/IOSPing.pm
smokeping/lib/probes/base.pm
smokeping/lib/probes/EchoPing.pm
smokeping/lib/probes/EchoPingChargen.pm
smokeping/lib/probes/EchoPingDiscard.pm
smokeping/lib/probes/EchoPingHttp.pm
smokeping/lib/probes/EchoPingHttps.pm
smokeping/lib/probes/EchoPingIcp.pm
smokeping/lib/probes/EchoPingSmtp.pm
smokeping/lib/probes/basefork.pm
smokeping/lib/probes/basevars.pm
smokeping/lib/ISG/ParseConfig.pm
@unexec rmdir %D/etc/smokeping || true
@exec mkdir -p %D/smokeping/htdocs/img; chown www:www %D/smokeping/htdocs/img
@unexec rm -rf %D/smokeping/htdocs/img
@exec mkdir -p %D/var/smokeping
@exec chown smokeping:smokeping %D/var/smokeping
@unexec rmdir %D/var/smokeping 2>/dev/null || true
@dirrm smokeping/htdocs
@dirrm smokeping/lib/probes
@dirrm smokeping/lib/ISG
@dirrm smokeping/lib
@dirrm share/doc/smokeping
@unexec rmdir %D/smokeping 2>/dev/null || true
