# net/mediaproxy/files/mediaproxy.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: mediaproxy
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# mediaproxy   starts and stops the SER MediaProxy server
#
# Add the following line to /etc/rc.conf to enable mediaproxy:
# mediaproxy_enable (bool):	Set to "NO" by default.
#				Set it to "YES" to enable mediaproxy.
#
# mediaproxy_flags (string):	Set command line options

. /etc/rc.subr

name="mediaproxy"
rcvar="${name}_enable"

load_rc_config $name

: ${mediaproxy_enable="NO"}

pidfile="/var/run/${name}.pid"

command="%%PREFIX%%/mediaproxy/mediaproxy.py"
command_interpreter="%%PYTHON_CMD%%"
command_args="--pid $pidfile"

start_postcmd="%%PREFIX%%/etc/rc.d/proxydispatcher onestart"
stop_postcmd="%%PREFIX%%/etc/rc.d/proxydispatcher stop"

run_rc_command "$1"
