# net/mediaproxy/files/proxydispatcher.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: proxydispatcher
# REQUIRE: mediaproxy
# KEYWORD: nostart
#
# proxydispatcher   starts and stops the SER MediaProxy dispatcher
#
# This script is usually run by %%PREFIX%%/etc/rc.d/mediaproxy
#
# Add the following line to /etc/rc.conf to enable proxydispatcher:
# proxydispatcher_enable (bool):	Set to "NO" by default.
#					Set it to "YES" to enable dispatcher.
# proxydispatcher_flags (string):	Set command line options

. /etc/rc.subr

name="proxydispatcher"
rcvar="${name}_enable"

load_rc_config $name

: ${proxydispatcher_enable="NO"}

pidfile="/var/run/${name}.pid"

command="%%PREFIX%%/mediaproxy/proxydispatcher.py"
command_interpreter="%%PYTHON_CMD%%"
command_args="--pid $pidfile"

run_rc_command "$1"
