# japanese/tiarra/files/tiarra.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: tiarra
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="tiarra"
rcvar=tiarra_enable

load_rc_config $name

tiarra_enable=${tiarra_enable:-"NO"}
tiarra_flags=${tiarra_flags:-"--quiet"}
tiarra_user=${tiarra_user:-"%%TIARRA_USER%%"}
tiarra_chdir=${tiarra_chdir:-"%%PREFIX%%/${name}"}

command="${tiarra_chdir}/${name}"
command_interpreter="/usr/bin/perl"
required_dirs="${tiarra_chdir}"
required_files="${tiarra_chdir}/${name}.conf"
extra_commands="reload"
sig_stop="${tiarra_sig_stop:-TERM}"

run_rc_command "$1"
