# security/openvas-server/files/openvasd.in
#!/bin/sh

# Start or stop openvasd
# $FreeBSD$

# PROVIDE: openvasd
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown
#

prefix=%%PREFIX%%

# Define these openvasd_* variables in one of these files:
#	/etc/rc.conf
#	/etc/rc.conf.local
#	/etc/rc.conf.d/openvasd
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#
[ -z "$openvasd_enable" ] && openvasd_enable="NO"	# Disable openvasd
#openvasd_program="${prefix}/sbin/openvasd"	# Location of openvasd
openvasd_flags="-D -q"				# Flags to openvasd program

. /etc/rc.subr

name="openvasd"
rcvar=`set_rcvar`
command="${prefix}/sbin/openvasd"
pidfile="/var/run/openvasd.pid"
required_files="${prefix}/etc/openvas/openvasd.conf"

load_rc_config $name
run_rc_command "$1"
