# sysutils/swapd/files/openrc-swapd.in
#!/sbin/openrc-run

name=swapd

command=%%PREFIX%%/sbin/swapd

depend() {
	keyword -shutdown
}
