# sysutils/daedalus/pkg-plist
bin/daedalus.rb
etc/rc.d/daedalus.sh
etc/daedalus.conf.sample
