# net/mDNSResponder/files/mdnsd.in
#!/bin/sh
#
# $FreeBSD$

# PROVIDE: mdnsd
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name=mdnsd
rcvar=`set_rcvar`

load_rc_config $name

: ${mdnsd_enable="NO"}
: ${mdnsd_pidfile="/var/run/${name}.pid"}

command="%%PREFIX%%/sbin/${name}"
pidfile="${mdnsd_pidfile}"

run_rc_command $*
