# net/mDNSResponder/files/mdnsresponderposix.in
#!/bin/sh

# PROVIDE: mdnsresponderposix
# REQUIRE: DAEMON
# KEYWORD: resume shutdown

. /etc/rc.subr

name=mdnsresponderposix
rcvar=mdnsresponderposix_enable

load_rc_config $name

: ${mdnsresponderposix_enable="NO"}
: ${mdnsresponderposix_pidfile="/var/run/${name}.pid"}

command="%%PREFIX%%/bin/mDNSResponderPosix"
command_args="-b -P ${mdnsresponderposix_pidfile}"
extra_commands="resume"
resume_cmd="mdnsresponderposix_resume"

mdnsresponderposix_resume() {
	run_rc_command restart
}

run_rc_command $*
