# mail/dk-milter/files/openrc-milter-dk.in
#!/sbin/openrc-run

# milterdk_uid (str):        Set username to run milter.
# milterdk_profiles (list):  Set to "" by default.
#                             Define your profiles here.
# milterdk_socket (str):    Path to the milter socket.
# milterdk_domain (str):    Domainpart of From: in mails to sign.
# milterdk_key (str):       Path to the private key file to sign with.
# milterdk_flags (str):     Flags passed to start command.
# milterdk_selector (str):  Selector from which the receiving side should
#			    obtain the public key.
#
# milterdk_${profile}_* :   Variables per profile.
#                             Sockets must be different from each other.

name="milterdk"

milterdk_uid=${milterdk_uid:-"mailnull"}
milterdk_profiles=${milterdk_profiles:-}
milterdk_socket=${milterdk_socket:-"unix:/var/run/milterdk/filter"}
milterdk_domain=${milterdk_domain:-"example.com"}
milterdk_key=${milterdk_key:-"/var/db/domainkeys/default.key.pem"}
milterdk_selector=${milterdk_selector:-"default"}
## milterdk_flags expands escaped variables later.
milterdk_flags=${milterdk_flags:-'-d ${milterdk_domain} -c nofws -H -m MSA -s ${milterdk_key} -S ${milterdk_selector}'}
command_args=milterdk_flags
command_user=mailnull

command="%%PREFIX%%/libexec/dk-filter"
_piddir="/var/run/milterdk"
pidfile="${_piddir}/pid"

if [ -n "$2" ]; then
    profile="$2"
    if [ "x${milterdk_profiles}" != "x" ]; then
	pidfile="${_piddir}/${profile}.pid"
	eval milterdk_enable="\${milterdk_${profile}_enable:-${milterdk_enable}}"
	eval milterdk_socket="\${milterdk_${profile}_socket:-}"
	if [ "x${milterdk_socket}" = "x" ];then
	    eerror "You must define a socket (milterdk_${profile}_socket)"
	    exit 1
	fi
	eval milterdk_domain="\${milterdk_${profile}_domain:-${milterdk_domain}}"
	eval milterdk_key="\${milterdk_${profile}_key:-${milterdk_key}}"
	eval milterdk_flags="\${milterdk_${profile}_flags:-${milterdk_flags}}"
	command_args="-l -p ${milterdk_socket} -u ${milterdk_uid} -P ${pidfile}"
    else
	einfo "$0: extra argument ignored"
    fi
else
    if [ "x${milterdk_profiles}" != "x" -a "x$1" != "x" ]; then
	if [ "x$1" != "xrestart" ]; then
	    for profile in ${milterdk_profiles}; do
		einfo "===> milterdk profile: ${profile}"
		%%RC_SCRIPT%% $1 ${profile}
		retcode="$?"
		if [ "0${retcode}" -ne 0 ]; then
		    failed="${profile} (${retcode}) ${failed:-}"
		else
		    success="${profile} ${success:-}"
		fi
	    done
	    exit 0
	else
	    restart_precmd=""
	fi
    else
	milterdk_flags=${milterdk_flags}
	command_args="-l -p ${milterdk_socket} -u ${milterdk_uid} -P ${pidfile}"
    fi
fi

start_pre()
{
    if [ -S ${milterdk_socket##local:} ] ; then
        rm -f ${milterdk_socket##local:}
    elif [ -S ${milterdk_socket##unix:} ] ; then
        rm -f ${milterdk_socket##unix:}
    fi
    if [ ! -d ${_piddir} ] ; then
	mkdir -p ${_piddir}
    fi
    if [ -n "${milterdk_uid}" ] ; then
	chown ${milterdk_uid} ${_piddir}
    fi
}

start_post()
{
    if [ -S ${milterdk_socket##local:} ] ; then
        rm -f ${milterdk_socket##local:}
    elif [ -S ${milterdk_socket##unix:} ] ; then
        rm -f ${milterdk_socket##unix:}
    fi
    # just if the directory is empty
    rmdir ${_piddir} > /dev/null 2>&1
}

depend()
{
	keyword -shutdown
}
