# net/wired/files/wired.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: wired
# REQUIRE: 
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# wired_enable (bool):   Set to NO by default.
#               Set it to YES to enable wired.
#

. /etc/rc.subr

name="wired"
rcvar=wired_enable

pidfile=%%PREFIX%%/wired/wired.pid
command=%%PREFIX%%/sbin/${name}

load_rc_config $name

: ${wired_enable="NO"}

run_rc_command "$1"
