# dns/dnsdist/files/dnsdist.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: dnsdist
# REQUIRE: NETWORKING DAEMON
# KEYWORD: shutdown
#
# Add the following to /etc/rc.conf to enable dnsdist:
#
# dnsdist_enable="YES"
#

. /etc/rc.subr

name="dnsdist"
rcvar="dnsdist_enable"

load_rc_config ${name}

: ${dnsdist_enable:="NO"}
: ${dnsdist_pidfile:=/var/run/${name}/pid}

command=/usr/local/sbin/${name}
command_args="-l 127.0.0.1:53"

run_rc_command "$1"
