# dns/dnsdist/files/dnsdist.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: dnsdist
# REQUIRE: NETWORKING DAEMON
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable dnsdist:
#
# dnsdist_enable="YES"

. /etc/rc.subr

desc="dnsdist is a powerful, scriptable DNS loadbalancer"
name=dnsdist
rcvar=dnsdist_enable

load_rc_config ${name}

: ${dnsdist_priv_user:=_dnsdist}
: ${dnsdist_priv_group:=_dnsdist}
: ${dnsdist_enable:=NO}

command=/usr/local/sbin/${name}
command_args="-u ${dnsdist_priv_user} -g ${dnsdist_priv_group} -d"

run_rc_command "$1"
