# games/tetrinet-x/files/openrc-tetrinetx.in
#!/sbin/openrc-run

name=tetrinetx
command="%%DATADIR%%/tetrix.freebsd"

depend() {
	keyword -shutdown
}
