# sysutils/3dm/files/openrc-3dm2.in
#!/sbin/openrc-run

name=twdm2

command="%%PREFIX%%/sbin/3dm2"
required_files=%%PREFIX%%/etc/3dm2/3dm2.conf

start_pre() {
	umask 077
}

depend() {
	use net
	keyword -shutdown
}
