# comms/ncid/pkg-plist
bin/cidalias
bin/cidcall
bin/cidupdate
bin/email2ncid
bin/ncid
bin/ncid-email2ncid-setup
bin/ncid-mysql-setup
bin/ncid-setup
bin/ncid2ncid
bin/ncidutil
bin/obi2ncid
bin/rn2ncid
bin/wc2ncid
bin/wct
bin/xdmf2ncid
bin/yac2ncid
@sample %%ETCDIR%%/email2ncid.conf.sample
@sample %%ETCDIR%%/ncid.conf.sample
@sample %%ETCDIR%%/ncid2ncid.conf.sample
@sample %%ETCDIR%%/ncidd.alias.sample
@sample %%ETCDIR%%/ncidd.blacklist.sample
@sample %%ETCDIR%%/ncidd.conf.sample
@sample %%ETCDIR%%/ncidd.whitelist.sample
@sample %%ETCDIR%%/obi2ncid.conf.sample
@sample %%ETCDIR%%/rn2ncid.conf.sample
@sample %%ETCDIR%%/sip2ncid.conf.sample
@sample %%ETCDIR%%/wc2ncid.conf.sample
@sample %%ETCDIR%%/xdmf2ncid.conf.sample
@sample %%ETCDIR%%/yac2ncid.conf.sample
etc/rc.d/ncid-initmodem
etc/rc.d/ncid-mysql
etc/rc.d/ncid-mythtv
etc/rc.d/ncid-notify
etc/rc.d/ncid-page
etc/rc.d/ncid-samba
etc/rc.d/ncid-speak
etc/rc.d/ncid-yac
etc/rc.d/ncid2ncid
etc/rc.d/ncidd
etc/rc.d/obi2ncid
etc/rc.d/rn2ncid
etc/rc.d/sip2ncid
etc/rc.d/wc2ncid
etc/rc.d/xdmf2ncid
etc/rc.d/yac2ncid
man/man1/cidalias.1.gz
man/man1/cidcall.1.gz
man/man1/cidupdate.1.gz
man/man1/email2ncid.1.gz
man/man1/get-areacodes-list.1.gz
man/man1/get-fcc-list.1.gz
man/man1/hangup-calls.1.gz
man/man1/hangup-closed-skel.1.gz
man/man1/hangup-combo-skel.1.gz
man/man1/hangup-fakenum.1.gz
man/man1/hangup-fcc.1.gz
man/man1/hangup-message-skel.1.gz
man/man1/hangup-nohangup.1.gz
man/man1/hangup-skel.1.gz
man/man1/ncid-alert.1.gz
man/man1/ncid-email2ncid-setup.1.gz
man/man1/ncid-initmodem.1.gz
man/man1/ncid-kpopup.1.gz
man/man1/ncid-mysql.1.gz
man/man1/ncid-mythtv.1.gz
man/man1/ncid-notify.1.gz
man/man1/ncid-page.1.gz
man/man1/ncid-samba.1.gz
man/man1/ncid-setup.1.gz
man/man1/ncid-skel.1.gz
man/man1/ncid-speak.1.gz
man/man1/ncid-wakeup.1.gz
man/man1/ncid-yac.1.gz
man/man1/ncid-yearlog.1.gz
man/man1/ncid.1.gz
man/man1/ncid2ncid.1.gz
man/man1/ncidrotate.1.gz
man/man1/ncidutil.1.gz
man/man1/obi2ncid.1.gz
man/man1/rn2ncid.1.gz
man/man1/wc2ncid.1.gz
man/man1/wct.1.gz
man/man1/xdmf2ncid.1.gz
man/man1/yac2ncid.1.gz
man/man5/email2ncid.conf.5.gz
man/man5/ncid.conf.5.gz
man/man5/ncid2ncid.conf.5.gz
man/man5/ncidd.alias.5.gz
man/man5/ncidd.blacklist.5.gz
man/man5/ncidd.conf.5.gz
man/man5/ncidd.whitelist.5.gz
man/man5/ncidrotate.conf.5.gz
man/man5/obi2ncid.conf.5.gz
man/man5/rn2ncid.conf.5.gz
man/man5/rotatebysize.conf.5.gz
man/man5/sip2ncid.conf.5.gz
man/man5/wc2ncid.conf.5.gz
man/man5/xdmf2ncid.conf.5.gz
man/man5/yac2ncid.conf.5.gz
man/man7/ncid_extensions.7.gz
man/man7/ncid_gateways.7.gz
man/man7/ncid_modules.7.gz
man/man7/ncid_recordings.7.gz
man/man7/ncid_tools.7.gz
man/man8/ncid-mysql-setup.8.gz
man/man8/ncidd.8.gz
man/man8/sip2ncid.8.gz
sbin/ncidd
sbin/sip2ncid
%%PORTDOCS%%%%DOCSDIR%%/Alias.md
%%PORTDOCS%%%%DOCSDIR%%/CONTRIBUTORS.md
%%PORTDOCS%%%%DOCSDIR%%/Clients.md
%%PORTDOCS%%%%DOCSDIR%%/Devices.md
%%PORTDOCS%%%%DOCSDIR%%/Extensions.md
%%PORTDOCS%%%%DOCSDIR%%/FAQ.md
%%PORTDOCS%%%%DOCSDIR%%/FakeNum.md
%%PORTDOCS%%%%DOCSDIR%%/FCCData.md
%%PORTDOCS%%%%DOCSDIR%%/GPL.md
%%PORTDOCS%%%%DOCSDIR%%/Gateways.md
%%PORTDOCS%%%%DOCSDIR%%/GettingStarted.md
%%PORTDOCS%%%%DOCSDIR%%/Hangup.md
%%PORTDOCS%%%%DOCSDIR%%/INSTALL-Cygwin.md
%%PORTDOCS%%%%DOCSDIR%%/INSTALL-DEB.md
%%PORTDOCS%%%%DOCSDIR%%/INSTALL-Fedora.md
%%PORTDOCS%%%%DOCSDIR%%/INSTALL-FreeBSD.md
%%PORTDOCS%%%%DOCSDIR%%/INSTALL-Mac.md
%%PORTDOCS%%%%DOCSDIR%%/INSTALL-Redhat.md
%%PORTDOCS%%%%DOCSDIR%%/INSTALL-Win.md
%%PORTDOCS%%%%DOCSDIR%%/INSTALL.md
%%PORTDOCS%%%%DOCSDIR%%/InstallIndex.md
%%PORTDOCS%%%%DOCSDIR%%/Log.md
%%PORTDOCS%%%%DOCSDIR%%/Message.md
%%PORTDOCS%%%%DOCSDIR%%/Modems.md
%%PORTDOCS%%%%DOCSDIR%%/Modules.md
%%PORTDOCS%%%%DOCSDIR%%/NCID-API.md
%%PORTDOCS%%%%DOCSDIR%%/NCID-SDK.md
%%PORTDOCS%%%%DOCSDIR%%/NCID-UserManual.md
%%PORTDOCS%%%%DOCSDIR%%/Obtain.md
%%PORTDOCS%%%%DOCSDIR%%/README-FreeBSD
%%PORTDOCS%%%%DOCSDIR%%/ReleaseNotes.md
%%PORTDOCS%%%%DOCSDIR%%/TODO.md
%%PORTDOCS%%%%DOCSDIR%%/Tools.md
%%PORTDOCS%%%%DOCSDIR%%/Use.md
%%PORTDOCS%%%%DOCSDIR%%/UserManualIndex.md
%%PORTDOCS%%%%DOCSDIR%%/Verbose.md
%%PORTDOCS%%%%DOCSDIR%%/recordings/CallingDeposit.pvf
%%PORTDOCS%%%%DOCSDIR%%/recordings/CannotBeCompleted.pvf
%%PORTDOCS%%%%DOCSDIR%%/recordings/DisconnectedNotInService.pvf
%%PORTDOCS%%%%DOCSDIR%%/recordings/NotInService.pvf
%%PORTDOCS%%%%DOCSDIR%%/recordings/README-recordings
%%DATADIR%%/extensions/hangup-calls
%%DATADIR%%/extensions/hangup-closed-skel
%%DATADIR%%/extensions/hangup-combo-skel
%%DATADIR%%/extensions/hangup-fakenum
%%DATADIR%%/extensions/hangup-fcc
%%DATADIR%%/extensions/hangup-message-skel
%%DATADIR%%/extensions/hangup-nohangup
%%DATADIR%%/extensions/hangup-skel
%%DATADIR%%/modules/ncid-alert
%%DATADIR%%/modules/ncid-initmodem
%%DATADIR%%/modules/ncid-kpopup
%%DATADIR%%/modules/ncid-mysql
%%DATADIR%%/modules/ncid-mythtv
%%DATADIR%%/modules/ncid-notify
%%DATADIR%%/modules/ncid-page
%%DATADIR%%/modules/ncid-samba
%%DATADIR%%/modules/ncid-skel
%%DATADIR%%/modules/ncid-speak
%%DATADIR%%/modules/ncid-wakeup
%%DATADIR%%/modules/ncid-yac
%%DATADIR%%/recordings/CallingDeposit.rmd
%%DATADIR%%/recordings/CannotBeCompleted.rmd
%%DATADIR%%/recordings/DisconnectedNotInService.rmd
%%DATADIR%%/recordings/NotInService.rmd
share/pixmaps/ncid
%%LOGROTATE%%%%DATADIR%%/sys/ncid-yearlog
%%LOGROTATE%%%%DATADIR%%/sys/ncidrotate
%%LOGROTATE%%@sample %%ETCDIR%%/ncidrotate.conf.sample
%%LOGROTATE%%etc/logrotate.d/ncid
%%LOGROTATE%%man/man1/ncidrotate.1.gz
%%LOGROTATE%%man/man5/ncidrotate.conf.5.gz
%%LOGROTATE%%man/man5/rotatebysize.conf.5.gz
