# comms/ncid/pkg-plist
bin/cidalias
bin/cidcall
bin/cidupdate
bin/ncid
bin/ncid2ncid
bin/ncidutil
bin/rn2ncid
bin/wc2ncid
bin/wct
bin/yac2ncid
%%ETCDIR%%/ncidd.alias
etc/rc.d/ncid-initmodem
etc/rc.d/ncid-mythtv
etc/rc.d/ncid-notify
etc/rc.d/ncid-page
etc/rc.d/ncid-samba
etc/rc.d/ncid-speak
etc/rc.d/ncid-yac
etc/rc.d/ncid2ncid
etc/rc.d/ncidd
etc/rc.d/rn2ncid
etc/rc.d/sip2ncid
etc/rc.d/wc2ncid
etc/rc.d/yac2ncid
man/man1/cidalias.1.gz
man/man1/cidcall.1.gz
man/man1/cidupdate.1.gz
man/man1/ncid-alert.1.gz
man/man1/ncid-initmodem.1.gz
man/man1/ncid-kpopup.1.gz
man/man1/ncid-mythtv.1.gz
man/man1/ncid-notify.1.gz
man/man1/ncid-page.1.gz
man/man1/ncid-samba.1.gz
man/man1/ncid-skel.1.gz
man/man1/ncid-speak.1.gz
man/man1/ncid-wakeup.1.gz
man/man1/ncid-yac.1.gz
man/man1/ncid.1.gz
man/man1/ncid2ncid.1.gz
man/man1/ncidutil.1.gz
man/man1/obi2ncid.1.gz
man/man1/rn2ncid.1.gz
man/man1/wc2ncid.1.gz
man/man1/wct.1.gz
man/man1/yac2ncid.1.gz
man/man5/ncid.conf.5.gz
man/man5/ncid2ncid.conf.5.gz
man/man5/ncidd.alias.5.gz
man/man5/ncidd.blacklist.5.gz
man/man5/ncidd.conf.5.gz
man/man5/ncidd.whitelist.5.gz
man/man5/obi2ncid.conf.5.gz
man/man5/rn2ncid.conf.5.gz
man/man5/sip2ncid.conf.5.gz
man/man5/wc2ncid.conf.5.gz
man/man5/yac2ncid.conf.5.gz
man/man7/ncid-modules.7.gz
man/man7/ncidgateways.7.gz
man/man7/ncidtools.7.gz
man/man8/ncidd.8.gz
man/man8/sip2ncid.8.gz
sbin/ncidd
sbin/sip2ncid
%%DATADIR%%
share/pixmaps/ncid
