# comms/ncid/pkg-plist
sbin/ncidd
sbin/ncidsip
sbin/sip2ncid
bin/ncid
bin/cidalias
bin/cidcall
bin/cidupdate
bin/yac2ncid
etc/logrotate.d/ncid
@dirrm etc/logrotate.d
etc/rc.d/ncid-kpopup
etc/rc.d/ncid-page
etc/rc.d/ncid-samba
etc/rc.d/ncid-yac
etc/rc.d/ncidsip
etc/rc.d/yac2ncid
etc/rc.d/ncid-mythtv
etc/rc.d/ncid-speak
etc/rc.d/ncidd
etc/rc.d/sip2ncid
share/ncid/ncidrotate
share/ncid/ncid-kpopup
share/ncid/ncid-mythtv
share/ncid/ncid-page
share/ncid/ncid-samba
share/ncid/ncid-skel
share/ncid/ncid-speak
share/ncid/ncid-tivo
share/ncid/ncid-yac
@dirrm share/ncid
share/pixmaps/ncid.gif
