# emulators/vmware2/pkg-plist
bin/vmware
etc/rc.d/vmware.sh
etc/vmware/config
lib/vmware/bin/df
lib/vmware/bin/vmware
lib/vmware/bin/vmnet-bridge
lib/vmware/bin/vmnet-dhcpd
lib/vmware/bin/vmnet-sniffer
lib/vmware/bin/vmware-config.pl
lib/vmware/bin/vmware-mount.pl
lib/vmware/bin/vmware-uninstall.pl
lib/vmware/bin/vmware-loop
lib/vmware/bin/vmware-ping
lib/vmware/bin/vmware-wizard
lib/vmware/lib/config
lib/vmware/lib/floppies/linux.flp
lib/vmware/lib/floppies/windows.flp
lib/vmware/lib/help/AutoFit.htm
lib/vmware/lib/help/ConfigDiskPartition.htm
lib/vmware/lib/help/ConfigDisks.htm
lib/vmware/lib/help/ConfigEditor.htm
lib/vmware/lib/help/ConfigEthernet.htm
lib/vmware/lib/help/ConfigFloppy.htm
lib/vmware/lib/help/ConfigGuestOS.htm
lib/vmware/lib/help/ConfigMemory.htm
lib/vmware/lib/help/ConfigMisc.htm
lib/vmware/lib/help/ConfigParallel.htm
lib/vmware/lib/help/ConfigScsiDisks.htm
lib/vmware/lib/help/ConfigSerial.htm
lib/vmware/lib/help/ConfigSound.htm
lib/vmware/lib/help/ConfigWiz.htm
lib/vmware/lib/help/Configuration.htm
lib/vmware/lib/help/Fit.htm
lib/vmware/lib/help/FullScreen.gif
lib/vmware/lib/help/Fullscreen.htm
lib/vmware/lib/help/GeneralInput.htm
lib/vmware/lib/help/GeneralMemory.htm
lib/vmware/lib/help/GeneralPreferences.htm
lib/vmware/lib/help/Glossary.htm
lib/vmware/lib/help/Help.gif
lib/vmware/lib/help/Interacting.htm
lib/vmware/lib/help/Introduction.htm
lib/vmware/lib/help/Power.htm
lib/vmware/lib/help/PowerOff.gif
lib/vmware/lib/help/PowerOff.htm
lib/vmware/lib/help/PowerOn.gif
lib/vmware/lib/help/PowerOn.htm
lib/vmware/lib/help/PowerReset.gif
lib/vmware/lib/help/PowerReset.htm
lib/vmware/lib/help/Resume.gif
lib/vmware/lib/help/Resume.htm
lib/vmware/lib/help/RunVM.htm
lib/vmware/lib/help/SendCtlAltDel.htm
lib/vmware/lib/help/Setremdev.htm
lib/vmware/lib/help/Settings.htm
lib/vmware/lib/help/ShowAllHints.htm
lib/vmware/lib/help/StatusBar.htm
lib/vmware/lib/help/Suspend.gif
lib/vmware/lib/help/Suspend.htm
lib/vmware/lib/help/ToolbarButtons.htm
lib/vmware/lib/help/ToolsInstall.htm
lib/vmware/lib/help/View.htm
lib/vmware/lib/help/Wizcd.htm
lib/vmware/lib/help/Wizcongrats.htm
lib/vmware/lib/help/Wizdisk.htm
lib/vmware/lib/help/Wizdisk3.htm
lib/vmware/lib/help/WizdiskType.htm
lib/vmware/lib/help/Wizfloppy.htm
lib/vmware/lib/help/Wiznetwork.htm
lib/vmware/lib/help/Wizos.htm
lib/vmware/lib/help/Wizpath.htm
lib/vmware/lib/help/toc.htm
%%FREEBSD5%%lib/vmware/lib/modules/linker.hints
lib/vmware/lib/modules/vmmon_smp.ko
lib/vmware/lib/modules/vmmon_up.ko
lib/vmware/lib/xkeymap/be101
lib/vmware/lib/xkeymap/be104
lib/vmware/lib/xkeymap/chde101
lib/vmware/lib/xkeymap/chde104
lib/vmware/lib/xkeymap/chfr101
lib/vmware/lib/xkeymap/chfr104
lib/vmware/lib/xkeymap/de101
lib/vmware/lib/xkeymap/de104
lib/vmware/lib/xkeymap/dk101
lib/vmware/lib/xkeymap/dk104
lib/vmware/lib/xkeymap/es101
lib/vmware/lib/xkeymap/es104
lib/vmware/lib/xkeymap/fi101
lib/vmware/lib/xkeymap/fi104
lib/vmware/lib/xkeymap/fr101
lib/vmware/lib/xkeymap/fr104
lib/vmware/lib/xkeymap/gb101
lib/vmware/lib/xkeymap/gb104
lib/vmware/lib/xkeymap/is101
lib/vmware/lib/xkeymap/is104
lib/vmware/lib/xkeymap/it101
lib/vmware/lib/xkeymap/it104
lib/vmware/lib/xkeymap/jp106
lib/vmware/lib/xkeymap/jp109
lib/vmware/lib/xkeymap/no101
lib/vmware/lib/xkeymap/no104
lib/vmware/lib/xkeymap/pt101
lib/vmware/lib/xkeymap/pt104
lib/vmware/lib/xkeymap/se101
lib/vmware/lib/xkeymap/se104
lib/vmware/lib/xkeymap/us101
lib/vmware/lib/xkeymap/us104
share/doc/vmware/CHANGES
share/doc/vmware/DHCP-COPYRIGHT
share/doc/vmware/EULA
share/doc/vmware/Hints.FreeBSD
share/doc/vmware/INSTALL
share/doc/vmware/README
share/doc/vmware/README.FreeBSD
share/doc/vmware/SAMBA-LICENSE
@exec mkdir -p %%LINUXBASE%%/dev;for n in 0 1 2 3 4 5 6 7 8 9; do ln -s /dev/ttyv$n %%LINUXBASE%%/dev/tty$n;done
@exec ln -s /dev/ttyva %%LINUXBASE%%/dev/tty10
@exec ln -s /dev/ttyvb %%LINUXBASE%%/dev/tty11
@exec mknod %%LINUXBASE%%/dev/null c 2 2
@exec chmod 666 %%LINUXBASE%%/dev/null
@exec mknod %%LINUXBASE%%/dev/hda b 0 0x00010002
@exec mknod %%LINUXBASE%%/dev/hdb b 0 0x0001000a
@exec ln -sf %D/etc/vmware /etc/
@exec mknod %%LINUXBASE%%/dev/vmmon c 200 0
@exec mknod %%LINUXBASE%%/dev/vmnet1 c 149 %%VMNET1_MINOR%%
@unexec rm -f /etc/vmware %%LINUXBASE%%/dev/vmmon %%LINUXBASE%%/dev/vmnet1 %%LINUXBASE%%/bin/df %%LINUXBASE%%/dev/tty[0-9]* %%LINUXBASE%%/dev/hd? %%LINUXBASE%%/dev/null
@dirrm etc/vmware
@dirrm lib/vmware/bin
@dirrm lib/vmware/lib/floppies
@dirrm lib/vmware/lib/help
@dirrm lib/vmware/lib/modules
@dirrm lib/vmware/lib/xkeymap
@dirrm lib/vmware/lib
@dirrm lib/vmware
@dirrm share/doc/vmware
