# net/mrouted/files/mrouted.in
#!/bin/sh
#
# $FreeBSD$

# PROVIDE: mrouted
# REQUIRE: netif routing cleanvar
# KEYWORD: nojail

. /etc/rc.subr

name="mrouted"
rcvar="${name}_enable"
command="%%PREFIX%%/sbin/$name"
pidfile="/var/run/${name}.pid"
required_files="%%PREFIX%%/etc/${name}.conf"
extra_commands="reload"

load_rc_config $name
: ${mrouted_enable="NO"}
run_rc_command "$1"
