# net/pimd/files/pimd.in
#!/bin/sh
#
# PROVIDE: pimd
# REQUIRE: netif routing
# KEYWORD: nojail

#
# Add the following to /etc/rc.conf[.local] to enable this service
#
# pimd_enable="YES"
#

. /etc/rc.subr

name="pimd"
rcvar="${name}_enable"
command="%%PREFIX%%/sbin/${name}"
required_files="%%PREFIX%%/etc/${name}.conf"
extra_commands=reload

load_rc_config $name
: ${pimd_enable="NO"}
run_rc_command "$1"

