# net/wired-tracker/files/trackerd.sh.in
#!/bin/sh

# PROVIDE: trackerd
# REQUIRE: 
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# trackerd_enable (bool):   Set to NO by default.
#               Set it to YES to enable trackerd.
#

. /etc/rc.subr

name="trackerd"
rcvar=${name}_enable

pidfile=%%PREFIX%%/wired/trackerd.pid
command=%%PREFIX%%/sbin/${name}

load_rc_config $name

: ${trackerd_enable="NO"}

command_args=""

run_rc_command "$1"
