# mail/sendmail/pkg-plist
@exec if ! pw groupshow smmsp 2>/dev/null; then pw groupadd smmsp -g 25; fi
@exec if ! pw usershow smmsp 2>/dev/null; then pw useradd smmsp -g smmsp -u 25 -h - -d /nonexistent -s /nonexistent -c "Sendmail Queue"; fi
@exec mkdir -p /var/spool/clientmqueue
@exec chown smmsp:smmsp /var/spool/clientmqueue
@exec chmod 770 /var/spool/clientmqueue
@unexec if pw usershow smmsp | grep -q 90:90; then pw userdel smmsp; fi
@comment (removed by userdel) pw groupdel smmsp
bin/hoststat
bin/mailq
bin/newaliases
bin/purgestat
bin/rmail
bin/vacation
%%ETCPORT%%etc/rc.d/sm-client.sh.sample
%%ETCPORT%%etc/rc.d/sendmail.sh.sample
libexec/mail.local
libexec/smrsh
sbin/editmap
sbin/mailstats
sbin/makemap
sbin/praliases
@group smmsp
sbin/sendmail
@group
@comment (overwritten by pkg_add) chown root:smmsp sbin/sendmail
@exec chmod 2555 %D/sbin/sendmail
share/sendmail/helpfile
%%PORTDOCS%%%%DOCSDIR%%/op.ps
%%PORTDOCS%%%%DOCSDIR%%/op.txt
%%PORTDOCS%%%%DOCSDIR%%/DEVTOOLS
%%PORTDOCS%%%%DOCSDIR%%/SENDMAIL
%%PORTDOCS%%%%DOCSDIR%%/MAIL.LOCAL
%%PORTDOCS%%%%DOCSDIR%%/SMRSH
