# mail/sendmail/pkg-plist
@comment $FreeBSD: /tmp/pcvs/ports/mail/sendmail/pkg-plist,v 1.10 2002-01-05 23:43:12 dinoex Exp $
@exec if ! pw groupshow smmsp 2>/dev/null; then pw groupadd smmsp -g 90; fi
@exec if ! pw usershow smmsp 2>/dev/null; then pw useradd smmsp -g smmsp -u 90 -h - -d /nonexistent -s /nonexistent -c "Sendmail Queue"; fi
@exec mkdir -p /var/spool/clientmqueue
@exec chown smmsp:smmsp /var/spool/clientmqueue
@exec chmod 770 /var/spool/clientmqueue
@unexec if pw usershow smmsp | grep -q 90:90; then pw userdel smmsp; fi
@comment (removed by userdel) pw groupdel smmsp
bin/hoststat
bin/mailq
bin/newaliases
bin/purgestat
bin/rmail
bin/vacation
etc/rc.d/sm-client.sh.sample
libexec/mail.local
libexec/smrsh
sbin/editmap
sbin/mailstats
sbin/makemap
sbin/praliases
@group smmsp
sbin/sendmail
@group
@comment (overwritten by pkg_add) chown root:smmsp sbin/sendmail
@exec chmod 2555 sbin/sendmail
%%PORTDOCS%%share/doc/sendmail/op.ps
%%PORTDOCS%%share/doc/sendmail/DEVTOOLS
%%PORTDOCS%%share/doc/sendmail/SENDMAIL
%%PORTDOCS%%share/doc/sendmail/MAIL.LOCAL
%%PORTDOCS%%share/doc/sendmail/SMRSH
