# net/netatalk-asun/pkg-plist
bin/adv1tov2
bin/aecho
bin/getzones
bin/hqx2bin
bin/macbinary
bin/megatron
bin/nbplkup
bin/nbprgstr
bin/nbpunrgstr
bin/pap
bin/papstatus
bin/psorder
bin/single2bin
bin/unbin
bin/unhex
bin/unsingle
@unexec if cmp -s %D/etc/AppleVolumes.default %D/etc/AppleVolumes.default.dist; then rm -f %D/etc/AppleVolumes.default; fi
etc/AppleVolumes.default.dist
@unexec if cmp -s %D/etc/AppleVolumes.system %D/etc/AppleVolumes.system.dist; then rm -f %D/etc/AppleVolumes.system; fi
etc/AppleVolumes.system.dist
@unexec if cmp -s %D/etc/afpd.conf %D/etc/afpd.conf.dist; then rm -f %D/etc/afpd.conf; fi
etc/afpd.conf.dist
etc/rc.d/netatalk.sh
include/atalk/adouble.h
include/atalk/aep.h
include/atalk/afp.h
include/atalk/asp.h
include/atalk/atp.h
include/atalk/cnid.h
include/atalk/compat.h
include/atalk/ddp.h
include/atalk/dsi.h
include/atalk/nbp.h
include/atalk/pap.h
include/atalk/paths.h
include/atalk/rtmp.h
include/atalk/server_child.h
include/atalk/util.h
include/atalk/zip.h
lib/libatalk.a
lib/libatalk_p.a
libexec/afpd
libexec/atalkd
libexec/etc2ps
libexec/papd
libexec/psa
libexec/psf
libexec/filters/ifmpap
libexec/filters/ifmpaprev
libexec/filters/ifpap
libexec/filters/ifpaprev
libexec/filters/ifwmpap
libexec/filters/ifwmpaprev
libexec/filters/ifwpap
libexec/filters/ifwpaprev
libexec/filters/ofmpap
libexec/filters/ofpap
libexec/filters/ofwmpap
libexec/filters/ofwpap
libexec/filters/tfmpap
libexec/filters/tfmpaprev
libexec/filters/tfpap
libexec/filters/tfpaprev
libexec/filters/tfwmpap
libexec/filters/tfwmpaprev
libexec/filters/tfwpap
libexec/filters/tfwpaprev
share/examples/netatalk/atalkd.conf
share/examples/netatalk/papd.conf
share/nls/netatalk/maccode.437
share/nls/netatalk/maccode.850
share/nls/netatalk/maccode.iso8859-1
share/pagecount.ps
@dirrm include/atalk
@dirrm libexec/filters
@dirrm share/examples/netatalk
@dirrm share/nls/netatalk
