# mail/sqlgrey/files/sqlgrey.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: sqlgrey
# REQUIRE: LOGIN
# BEFORE: mail
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable sqlgrey:
#
# sqlgrey_enable="YES"
#
# See man sqlgrey for flags or the config file.

. /etc/rc.subr

name=sqlgrey
rcvar=sqlgrey_enable

command=%%PREFIX%%/sbin/sqlgrey
extra_commands=reload
command_interpreter=%%PERL%%

stop_postcmd=${name}_poststop

load_rc_config $name

sqlgrey_enable=${sqlgrey_enable-"NO"}
required_files=${sqlgrey_config-"%%PREFIX%%/%%ETCDIR%%/sqlgrey.conf"}
pidfile=${sqlgrey_pidfile-"/var/run/sqlgrey.pid"}

command_args="--configfile=${required_files} --pidfile=${pidfile} --daemonize"

sqlgrey_poststop()
{
	rm -f $pidfile
}

run_rc_command "$1"
