# graphics/linux-c7-dri/pkg-plist.i386
usr/lib/clang-private/libclang.so
usr/lib/clang-private/libclang.so.7
usr/lib/clang-private/libclangARCMigrate.so
usr/lib/clang-private/libclangARCMigrate.so.7
usr/lib/clang-private/libclangAST.so
usr/lib/clang-private/libclangAST.so.7
usr/lib/clang-private/libclangASTMatchers.so
usr/lib/clang-private/libclangASTMatchers.so.7
usr/lib/clang-private/libclangAnalysis.so
usr/lib/clang-private/libclangAnalysis.so.7
usr/lib/clang-private/libclangBasic.so
usr/lib/clang-private/libclangBasic.so.7
usr/lib/clang-private/libclangCodeGen.so
usr/lib/clang-private/libclangCodeGen.so.7
usr/lib/clang-private/libclangCrossTU.so
usr/lib/clang-private/libclangCrossTU.so.7
usr/lib/clang-private/libclangDriver.so
usr/lib/clang-private/libclangDriver.so.7
usr/lib/clang-private/libclangDynamicASTMatchers.so
usr/lib/clang-private/libclangDynamicASTMatchers.so.7
usr/lib/clang-private/libclangEdit.so
usr/lib/clang-private/libclangEdit.so.7
usr/lib/clang-private/libclangFormat.so
usr/lib/clang-private/libclangFormat.so.7
usr/lib/clang-private/libclangFrontend.so
usr/lib/clang-private/libclangFrontend.so.7
usr/lib/clang-private/libclangFrontendTool.so
usr/lib/clang-private/libclangFrontendTool.so.7
usr/lib/clang-private/libclangHandleCXX.so
usr/lib/clang-private/libclangHandleCXX.so.7
usr/lib/clang-private/libclangHandleLLVM.so
usr/lib/clang-private/libclangHandleLLVM.so.7
usr/lib/clang-private/libclangIndex.so
usr/lib/clang-private/libclangIndex.so.7
usr/lib/clang-private/libclangLex.so
usr/lib/clang-private/libclangLex.so.7
usr/lib/clang-private/libclangParse.so
usr/lib/clang-private/libclangParse.so.7
usr/lib/clang-private/libclangRewrite.so
usr/lib/clang-private/libclangRewrite.so.7
usr/lib/clang-private/libclangRewriteFrontend.so
usr/lib/clang-private/libclangRewriteFrontend.so.7
usr/lib/clang-private/libclangSema.so
usr/lib/clang-private/libclangSema.so.7
usr/lib/clang-private/libclangSerialization.so
usr/lib/clang-private/libclangSerialization.so.7
usr/lib/clang-private/libclangStaticAnalyzerCheckers.so
usr/lib/clang-private/libclangStaticAnalyzerCheckers.so.7
usr/lib/clang-private/libclangStaticAnalyzerCore.so
usr/lib/clang-private/libclangStaticAnalyzerCore.so.7
usr/lib/clang-private/libclangStaticAnalyzerFrontend.so
usr/lib/clang-private/libclangStaticAnalyzerFrontend.so.7
usr/lib/clang-private/libclangTooling.so
usr/lib/clang-private/libclangTooling.so.7
usr/lib/clang-private/libclangToolingASTDiff.so
usr/lib/clang-private/libclangToolingASTDiff.so.7
usr/lib/clang-private/libclangToolingCore.so
usr/lib/clang-private/libclangToolingCore.so.7
usr/lib/clang-private/libclangToolingInclusions.so
usr/lib/clang-private/libclangToolingInclusions.so.7
usr/lib/clang-private/libclangToolingRefactor.so
usr/lib/clang-private/libclangToolingRefactor.so.7
usr/lib/dri/i915_dri.so
usr/lib/dri/i965_dri.so
usr/lib/dri/kms_swrast_dri.so
usr/lib/dri/nouveau_dri.so
usr/lib/dri/r300_dri.so
usr/lib/dri/r600_dri.so
usr/lib/dri/radeonsi_dri.so
usr/lib/dri/swrast_dri.so
usr/lib/dri/virtio_gpu_dri.so
usr/lib/dri/vmwgfx_dri.so
usr/lib/libEGL_mesa.so.0
usr/lib/libEGL_mesa.so.0.0.0
usr/lib/libGLU.so.1
usr/lib/libGLU.so.1.3.1
usr/lib/libGLX_mesa.so.0
usr/lib/libGLX_mesa.so.0.0.0
usr/lib/libGLX_system.so.0
usr/lib/libLLVM-7-rhel.so
usr/lib/libOSMesa.so.8
usr/lib/libOSMesa.so.8.0.0
usr/lib/libgbm.so.1
usr/lib/libgbm.so.1.0.0
usr/lib/libglapi.so.0
usr/lib/libglapi.so.0.0.0
usr/lib/libglut.so.3
usr/lib/libglut.so.3.10.0
usr/lib/libxatracker.so.2
usr/lib/libxatracker.so.2.4.0
usr/lib/vdpau/libvdpau_nouveau.so.1
usr/lib/vdpau/libvdpau_nouveau.so.1.0
usr/lib/vdpau/libvdpau_nouveau.so.1.0.0
usr/lib/vdpau/libvdpau_r600.so.1
usr/lib/vdpau/libvdpau_r600.so.1.0
usr/lib/vdpau/libvdpau_r600.so.1.0.0
usr/lib/vdpau/libvdpau_radeonsi.so.1
usr/lib/vdpau/libvdpau_radeonsi.so.1.0
usr/lib/vdpau/libvdpau_radeonsi.so.1.0.0
%%PORTDOCS%%usr/share/doc/freeglut-3.0.0/AUTHORS
%%PORTDOCS%%usr/share/doc/freeglut-3.0.0/COPYING
%%PORTDOCS%%usr/share/doc/freeglut-3.0.0/ChangeLog
%%PORTDOCS%%usr/share/doc/freeglut-3.0.0/README
%%PORTDOCS%%usr/share/doc/freeglut-3.0.0/download.html
%%PORTDOCS%%usr/share/doc/freeglut-3.0.0/freeglut.html
%%PORTDOCS%%usr/share/doc/freeglut-3.0.0/freeglut_logo.png
%%PORTDOCS%%usr/share/doc/freeglut-3.0.0/freeglut_user_interface.html
%%PORTDOCS%%usr/share/doc/freeglut-3.0.0/index.html
%%PORTDOCS%%usr/share/doc/freeglut-3.0.0/ogl_sm.png
%%PORTDOCS%%usr/share/doc/freeglut-3.0.0/progress.html
%%PORTDOCS%%usr/share/doc/freeglut-3.0.0/structure.html
%%PORTDOCS%%usr/share/doc/llvm-private-7.0.1/LICENSE.TXT
usr/share/drirc.d/00-mesa-defaults.conf
usr/share/glvnd/egl_vendor.d/50_mesa.json
