# textproc/xaira/files/xaira-server.in
#!/bin/sh
#

# Add the following line to /etc/rc.conf to enable xaira:
#
# xaira_enable="YES"

# PROVIDE: sarad

. /etc/rc.subr

name="xaira"
rcvar=`set_rcvar`

command=%%PREFIX%%/sbin/xaira_daemon

load_rc_config $name
run_rc_command "$1"
