# www/anyterm/files/anytermd.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: anytermd
# REQUIRE: NETWORKING DAEMON LOGIN cleanvar devfs
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable anytermd:
#
# anytermd_enable (bool):   Set to "NO" by default.
#                           Set it to "YES" to enable anytermd.
#
# anytermd_flags (str):     Set to "-c /usr/bin/login -u nobody --local-only"
#                           by default. See anytermd(1) for flags.
#
# anytermd_name (str):      Set to "anyterm" by default.
#

. /etc/rc.subr

name=anytermd
rcvar=`set_rcvar`

load_rc_config $name

# Set defaults
: ${anytermd_enable:="NO"}
: ${anytermd_flags:="-c /usr/bin/login -u nobody --local-only"}
: ${anytermd_name:="anyterm"}

pidfile=/var/run/$anytermd_name.pid
command=%%PREFIX%%/sbin/$name
command_args="$anytermd_flags --name $anytermd_name"

run_rc_command "$1"
