# security/openvas-scanner/files/openrc-openvassd.in
#!/sbin/openrc-run

name=openvassd
command="%%PREFIX%%/sbin/openvassd"
pidfile=/var/run/openvassd.pid

depend() {
		keyword -shutdown
}
