# net/freeradius2/pkg-plist
@unexec %D/etc/rc.d/radiusd.sh forcestop 2>/dev/null || true
bin/radclient
bin/radeapclient
bin/radlast
bin/radrelay
bin/radtest
bin/radwho
bin/radzap
bin/rlm_dbm_cat
bin/rlm_dbm_parser
bin/rlm_ippool_tool
bin/smbencrypt
etc/raddb/certs/demoCA/serial.old
etc/raddb/certs/demoCA/serial
etc/raddb/certs/demoCA/index.txt.old
etc/raddb/certs/demoCA/index.txt
etc/raddb/certs/demoCA/cacert.pem
etc/raddb/certs/root.p12
etc/raddb/certs/root.der
etc/raddb/certs/random
etc/raddb/certs/newreq.pem
etc/raddb/certs/newcert.pem
etc/raddb/certs/dh
etc/raddb/certs/cert-srv.pem
etc/raddb/certs/cert-srv.p12
etc/raddb/certs/cert-srv.der
etc/raddb/certs/cert-clt.pem
etc/raddb/certs/cert-clt.p12
etc/raddb/certs/cert-clt.der
etc/raddb/certs/root.pem
etc/raddb/certs/README
etc/raddb/acct_users.sample
etc/raddb/attrs.sample
etc/raddb/clients.sample
etc/raddb/clients.conf.sample
etc/raddb/dictionary.sample
%%EXPM%%etc/raddb/example.pl
etc/raddb/experimental.conf.sample
etc/raddb/eap.conf.sample
etc/raddb/hints.sample
etc/raddb/huntgroups.sample
etc/raddb/ldap.attrmap.sample
etc/raddb/mssql.conf.sample
etc/raddb/naslist.sample
etc/raddb/naspasswd.sample
etc/raddb/oraclesql.conf.sample
etc/raddb/postgresql.conf.sample
etc/raddb/preproxy_users.sample
etc/raddb/proxy.conf.sample
etc/raddb/radiusd.conf.sample
etc/raddb/realms.sample
etc/raddb/snmp.conf.sample
etc/raddb/sql.conf.sample
etc/raddb/users.sample
etc/raddb/x99.conf.sample
etc/raddb/x99passwd.sample
etc/rc.d/radiusd.sh
lib/libeap-%%PORTVERSION%%.la
lib/libeap-%%PORTVERSION%%.so
lib/libeap.a
lib/libeap.la
lib/libeap.so
lib/libradius-%%PORTVERSION%%.la
lib/libradius-%%PORTVERSION%%.so
lib/libradius.a
lib/libradius.la
lib/libradius.so
lib/rlm_acct_unique-%%PORTVERSION%%.la
lib/rlm_acct_unique-%%PORTVERSION%%.so
lib/rlm_acct_unique.a
lib/rlm_acct_unique.la
lib/rlm_acct_unique.so
lib/rlm_always-%%PORTVERSION%%.la
lib/rlm_always-%%PORTVERSION%%.so
lib/rlm_always.a
lib/rlm_always.la
lib/rlm_always.so
lib/rlm_attr_filter-%%PORTVERSION%%.la
lib/rlm_attr_filter-%%PORTVERSION%%.so
lib/rlm_attr_filter.a
lib/rlm_attr_filter.la
lib/rlm_attr_filter.so
lib/rlm_attr_rewrite-%%PORTVERSION%%.la
lib/rlm_attr_rewrite-%%PORTVERSION%%.so
lib/rlm_attr_rewrite.a
lib/rlm_attr_rewrite.la
lib/rlm_attr_rewrite.so
lib/rlm_chap-%%PORTVERSION%%.la
lib/rlm_chap-%%PORTVERSION%%.so
lib/rlm_chap.a
lib/rlm_chap.la
lib/rlm_chap.so
lib/rlm_checkval-%%PORTVERSION%%.la
lib/rlm_checkval-%%PORTVERSION%%.so
lib/rlm_checkval.a
lib/rlm_checkval.la
lib/rlm_checkval.so
lib/rlm_counter-%%PORTVERSION%%.la
lib/rlm_counter-%%PORTVERSION%%.so
lib/rlm_counter.a
lib/rlm_counter.la
lib/rlm_counter.so
%%EXPM%%lib/rlm_cram-%%PORTVERSION%%.la
%%EXPM%%lib/rlm_cram-%%PORTVERSION%%.so
%%EXPM%%lib/rlm_cram.a
%%EXPM%%lib/rlm_cram.la
%%EXPM%%lib/rlm_cram.so
lib/rlm_dbm-%%PORTVERSION%%.la
lib/rlm_dbm-%%PORTVERSION%%.so
lib/rlm_dbm.a
lib/rlm_dbm.la
lib/rlm_dbm.so
lib/rlm_detail-%%PORTVERSION%%.la
lib/rlm_detail-%%PORTVERSION%%.so
lib/rlm_detail.a
lib/rlm_detail.la
lib/rlm_detail.so
lib/rlm_digest-%%PORTVERSION%%.la
lib/rlm_digest-%%PORTVERSION%%.so
lib/rlm_digest.a
lib/rlm_digest.la
lib/rlm_digest.so
lib/rlm_eap-%%PORTVERSION%%.la
lib/rlm_eap-%%PORTVERSION%%.so
lib/rlm_eap.a
lib/rlm_eap.la
lib/rlm_eap.so
lib/rlm_eap_gtc-%%PORTVERSION%%.la
lib/rlm_eap_gtc-%%PORTVERSION%%.so
lib/rlm_eap_gtc.a
lib/rlm_eap_gtc.la
lib/rlm_eap_gtc.so
lib/rlm_eap_leap-%%PORTVERSION%%.la
lib/rlm_eap_leap-%%PORTVERSION%%.so
lib/rlm_eap_leap.a
lib/rlm_eap_leap.la
lib/rlm_eap_leap.so
lib/rlm_eap_md5-%%PORTVERSION%%.la
lib/rlm_eap_md5-%%PORTVERSION%%.so
lib/rlm_eap_md5.a
lib/rlm_eap_md5.la
lib/rlm_eap_md5.so
lib/rlm_eap_mschapv2-%%PORTVERSION%%.la
lib/rlm_eap_mschapv2-%%PORTVERSION%%.so
lib/rlm_eap_mschapv2.a
lib/rlm_eap_mschapv2.la
lib/rlm_eap_mschapv2.so
lib/rlm_eap_peap-%%PORTVERSION%%.la
lib/rlm_eap_peap-%%PORTVERSION%%.so
lib/rlm_eap_peap.a
lib/rlm_eap_peap.la
lib/rlm_eap_peap.so
lib/rlm_eap_sim-%%PORTVERSION%%.la
lib/rlm_eap_sim-%%PORTVERSION%%.so
lib/rlm_eap_sim.a
lib/rlm_eap_sim.la
lib/rlm_eap_sim.so
lib/rlm_eap_tls-%%PORTVERSION%%.la
lib/rlm_eap_tls-%%PORTVERSION%%.so
lib/rlm_eap_tls.a
lib/rlm_eap_tls.la
lib/rlm_eap_tls.so
lib/rlm_eap_ttls-%%PORTVERSION%%.la
lib/rlm_eap_ttls-%%PORTVERSION%%.so
lib/rlm_eap_ttls.a
lib/rlm_eap_ttls.la
lib/rlm_eap_ttls.so
%%EXPM%%lib/rlm_example-%%PORTVERSION%%.la
%%EXPM%%lib/rlm_example-%%PORTVERSION%%.so
%%EXPM%%lib/rlm_example.a
%%EXPM%%lib/rlm_example.la
%%EXPM%%lib/rlm_example.so
lib/rlm_exec-%%PORTVERSION%%.la
lib/rlm_exec-%%PORTVERSION%%.so
lib/rlm_exec.a
lib/rlm_exec.la
lib/rlm_exec.so
lib/rlm_expr-%%PORTVERSION%%.la
lib/rlm_expr-%%PORTVERSION%%.so
lib/rlm_expr.a
lib/rlm_expr.la
lib/rlm_expr.so
lib/rlm_fastusers-%%PORTVERSION%%.la
lib/rlm_fastusers-%%PORTVERSION%%.so
lib/rlm_fastusers.a
lib/rlm_fastusers.la
lib/rlm_fastusers.so
lib/rlm_files-%%PORTVERSION%%.la
lib/rlm_files-%%PORTVERSION%%.so
lib/rlm_files.a
lib/rlm_files.la
lib/rlm_files.so
lib/rlm_ippool-%%PORTVERSION%%.la
lib/rlm_ippool-%%PORTVERSION%%.so
lib/rlm_ippool.a
lib/rlm_ippool.la
lib/rlm_ippool.so
%%KRB5%%lib/rlm_krb5-%%PORTVERSION%%.la
%%KRB5%%lib/rlm_krb5-%%PORTVERSION%%.so
%%KRB5%%lib/rlm_krb5.a
%%KRB5%%lib/rlm_krb5.la
%%KRB5%%lib/rlm_krb5.so
%%LDAP%%lib/rlm_ldap-%%PORTVERSION%%.la
%%LDAP%%lib/rlm_ldap-%%PORTVERSION%%.so
%%LDAP%%lib/rlm_ldap.a
%%LDAP%%lib/rlm_ldap.la
%%LDAP%%lib/rlm_ldap.so
lib/rlm_mschap-%%PORTVERSION%%.la
lib/rlm_mschap-%%PORTVERSION%%.so
lib/rlm_mschap.a
lib/rlm_mschap.la
lib/rlm_mschap.so
lib/rlm_ns_mta_md5-%%PORTVERSION%%.la
lib/rlm_ns_mta_md5-%%PORTVERSION%%.so
lib/rlm_ns_mta_md5.a
lib/rlm_ns_mta_md5.la
lib/rlm_ns_mta_md5.so
lib/rlm_pam-%%PORTVERSION%%.la
lib/rlm_pam-%%PORTVERSION%%.so
lib/rlm_pam.a
lib/rlm_pam.la
lib/rlm_pam.so
lib/rlm_pap-%%PORTVERSION%%.la
lib/rlm_pap-%%PORTVERSION%%.so
lib/rlm_pap.a
lib/rlm_pap.la
lib/rlm_pap.so
lib/rlm_passwd-%%PORTVERSION%%.la
lib/rlm_passwd-%%PORTVERSION%%.so
lib/rlm_passwd.a
lib/rlm_passwd.la
lib/rlm_passwd.so
%%EXPM%%lib/rlm_perl-%%PORTVERSION%%.la
%%EXPM%%lib/rlm_perl-%%PORTVERSION%%.so
%%EXPM%%lib/rlm_perl.a
%%EXPM%%lib/rlm_perl.la
%%EXPM%%lib/rlm_perl.so
lib/rlm_preprocess-%%PORTVERSION%%.la
lib/rlm_preprocess-%%PORTVERSION%%.so
lib/rlm_preprocess.a
lib/rlm_preprocess.la
lib/rlm_preprocess.so
%%EXPM%%lib/rlm_python-%%PORTVERSION%%.la
%%EXPM%%lib/rlm_python-%%PORTVERSION%%.so
%%EXPM%%lib/rlm_python.a
%%EXPM%%lib/rlm_python.la
%%EXPM%%lib/rlm_python.so
lib/rlm_radutmp-%%PORTVERSION%%.la
lib/rlm_radutmp-%%PORTVERSION%%.so
lib/rlm_radutmp.a
lib/rlm_radutmp.la
lib/rlm_radutmp.so
lib/rlm_realm-%%PORTVERSION%%.la
lib/rlm_realm-%%PORTVERSION%%.so
lib/rlm_realm.a
lib/rlm_realm.la
lib/rlm_realm.so
%%EXPM%%lib/rlm_sim_files-%%PORTVERSION%%.la
%%EXPM%%lib/rlm_sim_files-%%PORTVERSION%%.so
%%EXPM%%lib/rlm_sim_files.a
%%EXPM%%lib/rlm_sim_files.la
%%EXPM%%lib/rlm_sim_files.so
%%EXPM%%lib/rlm_smb-%%PORTVERSION%%.la
%%EXPM%%lib/rlm_smb-%%PORTVERSION%%.so
%%EXPM%%lib/rlm_smb.a
%%EXPM%%lib/rlm_smb.la
%%EXPM%%lib/rlm_smb.so
lib/rlm_sql-%%PORTVERSION%%.la
lib/rlm_sql-%%PORTVERSION%%.so
lib/rlm_sql.a
lib/rlm_sql.la
lib/rlm_sql.so
%%EXPM%%lib/rlm_sqlcounter-%%PORTVERSION%%.la
%%EXPM%%lib/rlm_sqlcounter-%%PORTVERSION%%.so
%%EXPM%%lib/rlm_sqlcounter.a
%%EXPM%%lib/rlm_sqlcounter.la
%%EXPM%%lib/rlm_sqlcounter.so
%%MYSQL%%lib/rlm_sql_mysql-%%PORTVERSION%%.so
%%MYSQL%%lib/rlm_sql_mysql.a
%%MYSQL%%lib/rlm_sql_mysql.la
%%MYSQL%%lib/rlm_sql_mysql.so
%%PGSQL%%lib/rlm_sql_postgresql-%%PORTVERSION%%.so
%%PGSQL%%lib/rlm_sql_postgresql.a
%%PGSQL%%lib/rlm_sql_postgresql.la
%%PGSQL%%lib/rlm_sql_postgresql.so
lib/rlm_unix-%%PORTVERSION%%.la
lib/rlm_unix-%%PORTVERSION%%.so
lib/rlm_unix.a
lib/rlm_unix.la
lib/rlm_unix.so
%%TOKEN%%lib/rlm_x99_token-%%PORTVERSION%%.la
%%TOKEN%%lib/rlm_x99_token-%%PORTVERSION%%.so
%%TOKEN%%lib/rlm_x99_token.a
%%TOKEN%%lib/rlm_x99_token.la
%%TOKEN%%lib/rlm_x99_token.so
sbin/check-radiusd-config
sbin/checkrad
sbin/radiusd
sbin/radwatch
sbin/rc.radiusd
%%PORTDOCS%%%%DOCSDIR%%/Autz-Type
%%PORTDOCS%%%%DOCSDIR%%/DIFFS
%%PORTDOCS%%%%DOCSDIR%%/MACOSX
%%PORTDOCS%%%%DOCSDIR%%/OS2
%%PORTDOCS%%%%DOCSDIR%%/RADIUS-LDAP.schema
%%PORTDOCS%%%%DOCSDIR%%/RADIUS-LDAPv3.schema
%%PORTDOCS%%%%DOCSDIR%%/RADIUS-SQL.schema
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/Simultaneous-Use
%%PORTDOCS%%%%DOCSDIR%%/aaa.txt
%%PORTDOCS%%%%DOCSDIR%%/ascend
%%PORTDOCS%%%%DOCSDIR%%/bay
%%PORTDOCS%%%%DOCSDIR%%/bugs
%%PORTDOCS%%%%DOCSDIR%%/cisco
%%PORTDOCS%%%%DOCSDIR%%/coding-methods.txt
%%PORTDOCS%%%%DOCSDIR%%/configurable_failover
%%PORTDOCS%%%%DOCSDIR%%/duplicate-users
%%PORTDOCS%%%%DOCSDIR%%/module_interface
%%PORTDOCS%%%%DOCSDIR%%/performance-testing
%%PORTDOCS%%%%DOCSDIR%%/processing_users_file
%%PORTDOCS%%%%DOCSDIR%%/proxy
%%PORTDOCS%%%%DOCSDIR%%/radrelay
%%PORTDOCS%%%%DOCSDIR%%/rfc/attributes.html
%%PORTDOCS%%%%DOCSDIR%%/rfc/draft-kamath-pppext-eap-mschapv2-00.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/draft-sterman-aaa-sip-00.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/leap.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/pppext-eap-sim-12.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc1157.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc1227.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc1448.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc1901.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc1905.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2058.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2059.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2138.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2139.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2243.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2289.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2433.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2548.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2618.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2619.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2620.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2621.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2716.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2759.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2809.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2865.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2866.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2867.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2868.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2869.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2882.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc2924.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc3162.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc3575.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc3576.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc3579.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc/rfc3580.txt
%%PORTDOCS%%%%DOCSDIR%%/rlm_attr_filter
%%PORTDOCS%%%%DOCSDIR%%/rlm_dbm
%%PORTDOCS%%%%DOCSDIR%%/rlm_digest
%%PORTDOCS%%%%DOCSDIR%%/rlm_eap
%%PORTDOCS%%%%DOCSDIR%%/rlm_fastusers
%%PORTDOCS%%%%DOCSDIR%%/rlm_krb5
%%PORTDOCS%%%%DOCSDIR%%/rlm_ldap
%%PORTDOCS%%%%DOCSDIR%%/rlm_pam
%%PORTDOCS%%%%DOCSDIR%%/rlm_passwd
%%PORTDOCS%%%%DOCSDIR%%/rlm_python
%%PORTDOCS%%%%DOCSDIR%%/rlm_sim_triplets
%%PORTDOCS%%%%DOCSDIR%%/rlm_sql
%%PORTDOCS%%%%DOCSDIR%%/rlm_sqlcounter
%%PORTDOCS%%%%DOCSDIR%%/rlm_x99_token
%%PORTDOCS%%%%DOCSDIR%%/supervise-radiusd.txt
%%PORTDOCS%%%%DOCSDIR%%/tuning_guide
%%PORTDOCS%%%%DOCSDIR%%/variables.txt
%%DATADIR%%/dictionary
%%DATADIR%%/dictionary.3com
%%DATADIR%%/dictionary.3gpp
%%DATADIR%%/dictionary.3gpp2
%%DATADIR%%/dictionary.acc
%%DATADIR%%/dictionary.alcatel
%%DATADIR%%/dictionary.alteon
%%DATADIR%%/dictionary.altiga
%%DATADIR%%/dictionary.aptis
%%DATADIR%%/dictionary.ascend
%%DATADIR%%/dictionary.bay
%%DATADIR%%/dictionary.bintec
%%DATADIR%%/dictionary.bristol
%%DATADIR%%/dictionary.cabletron
%%DATADIR%%/dictionary.cisco
%%DATADIR%%/dictionary.cisco.bbsm
%%DATADIR%%/dictionary.cisco.vpn3000
%%DATADIR%%/dictionary.cisco.vpn5000
%%DATADIR%%/dictionary.colubris
%%DATADIR%%/dictionary.columbia_university
%%DATADIR%%/dictionary.compat
%%DATADIR%%/dictionary.erx
%%DATADIR%%/dictionary.extreme
%%DATADIR%%/dictionary.foundry
%%DATADIR%%/dictionary.freeradius
%%DATADIR%%/dictionary.gandalf
%%DATADIR%%/dictionary.garderos
%%DATADIR%%/dictionary.gemtek
%%DATADIR%%/dictionary.itk
%%DATADIR%%/dictionary.juniper
%%DATADIR%%/dictionary.karlnet
%%DATADIR%%/dictionary.livingston
%%DATADIR%%/dictionary.localweb
%%DATADIR%%/dictionary.merit
%%DATADIR%%/dictionary.microsoft
%%DATADIR%%/dictionary.mikrotik
%%DATADIR%%/dictionary.navini
%%DATADIR%%/dictionary.netscreen
%%DATADIR%%/dictionary.nokia
%%DATADIR%%/dictionary.nomadix
%%DATADIR%%/dictionary.propel
%%DATADIR%%/dictionary.quintum
%%DATADIR%%/dictionary.redback
%%DATADIR%%/dictionary.redcreek
%%DATADIR%%/dictionary.shasta
%%DATADIR%%/dictionary.shiva
%%DATADIR%%/dictionary.sonicwall
%%DATADIR%%/dictionary.springtide
%%DATADIR%%/dictionary.telebit
%%DATADIR%%/dictionary.trapeze
%%DATADIR%%/dictionary.tunnel
%%DATADIR%%/dictionary.unix
%%DATADIR%%/dictionary.usr
%%DATADIR%%/dictionary.valemount
%%DATADIR%%/dictionary.versanet
%%DATADIR%%/dictionary.wispr
%%DATADIR%%/dictionary.xedia
%%MYSQL%%%%EXAMPLESDIR%%/db_mysql.sql
%%PGSQL%%%%EXAMPLESDIR%%/db_postgresql.sql
@dirrm %%DATADIR%%
@dirrm %%EXAMPLESDIR%%
%%PORTDOCS%%@dirrm %%DOCSDIR%%/rfc
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm etc/raddb/certs/demoCA
@dirrm etc/raddb/certs
@exec mkdir -p /var/log/raddb
@exec chmod -R og= /var/log/raddb
@exec mkdir -p /var/run/radiusd
@unexec rm -fr /var/run/radiusd
@unexec rmdir %D/etc/raddb 2> /dev/null || echo "You should remove %D/etc/raddb if you don't need it anymore."
