# net/osrtspproxy/pkg-plist
sbin/rtspproxy
etc/rc.d/rtspproxy.sh.sample
