# sysutils/automounter/pkg-plist
sbin/automounter
etc/rc.d/automounter
etc/automounter.conf.sample
