# net/rsync/pkg-plist
@comment $FreeBSD$
bin/rsync
etc/rsyncd.conf.sample
etc/rc.d/rsyncd.sh.sample
%%PORTDOCS%%share/doc/rsync/README
%%PORTDOCS%%share/doc/rsync/COPYING
%%PORTDOCS%%share/doc/rsync/tech_report.tex
%%PORTDOCS%%@dirrm share/doc/rsync
