# ftp/wzdftpd/pkg-plist
@comment $FreeBSD$
bin/siteconfig
bin/siteuptime
bin/sitewho
bin/wzd-config
etc/rc.d/wzdftpd.sh
etc/wzd.pem
etc/file_ginfo.txt
etc/file_help.txt
etc/file_swho.txt
etc/file_user.txt
etc/file_group.txt
etc/file_rules.txt
etc/file_users.txt
etc/file_vfs.txt
etc/file_who.txt
etc/file_groups.txt
include/wzdftpd/libwzd-auth/wzd_auth.h
include/wzdftpd/libwzd-auth/wzd_base64.h
include/wzdftpd/libwzd-auth/wzd_crypt.h
include/wzdftpd/libwzd-auth/wzd_krb5.h
include/wzdftpd/libwzd-auth/wzd_md5.h
include/wzdftpd/libwzd-auth/wzd_md5crypt.h
include/wzdftpd/libwzd-auth/wzd_pam.h
include/wzdftpd/libwzd-auth/wzd_sha1.h
include/wzdftpd/libwzd-auth/wzd_tls.h
include/wzdftpd/libwzd-base/dlist.h
include/wzdftpd/libwzd-base/hash.h
include/wzdftpd/libwzd-base/list.h
include/wzdftpd/libwzd-base/stack.h
include/wzdftpd/libwzd-base/wzd_strlcat.h
include/wzdftpd/libwzd-base/wzd_strtok_r.h
include/wzdftpd/libwzd-core/arch/bsd.h
include/wzdftpd/libwzd-core/arch/win32.h
include/wzdftpd/libwzd-core/ls.h
include/wzdftpd/libwzd-core/wzd_ClientThread.h
include/wzdftpd/libwzd-core/wzd_action.h
include/wzdftpd/libwzd-core/wzd_all.h
include/wzdftpd/libwzd-core/wzd_backend.h
include/wzdftpd/libwzd-core/wzd_cache.h
include/wzdftpd/libwzd-core/wzd_commands.h
include/wzdftpd/libwzd-core/wzd_configfile.h
include/wzdftpd/libwzd-core/wzd_configloader.h
include/wzdftpd/libwzd-core/wzd_crc32.h
include/wzdftpd/libwzd-core/wzd_crontab.h
include/wzdftpd/libwzd-core/wzd_data.h
include/wzdftpd/libwzd-core/wzd_debug.h
include/wzdftpd/libwzd-core/wzd_dir.h
include/wzdftpd/libwzd-core/wzd_events.h
include/wzdftpd/libwzd-core/wzd_file.h
include/wzdftpd/libwzd-core/wzd_fs.h
include/wzdftpd/libwzd-core/wzd_hardlimits.h
include/wzdftpd/libwzd-core/wzd_ip.h
include/wzdftpd/libwzd-core/wzd_libmain.h
include/wzdftpd/libwzd-core/wzd_log.h
include/wzdftpd/libwzd-core/wzd_messages.h
include/wzdftpd/libwzd-core/wzd_misc.h
include/wzdftpd/libwzd-core/wzd_mod.h
include/wzdftpd/libwzd-core/wzd_mutex.h
include/wzdftpd/libwzd-core/wzd_perm.h
include/wzdftpd/libwzd-core/wzd_ratio.h
include/wzdftpd/libwzd-core/wzd_section.h
include/wzdftpd/libwzd-core/wzd_shm.h
include/wzdftpd/libwzd-core/wzd_site.h
include/wzdftpd/libwzd-core/wzd_site_group.h
include/wzdftpd/libwzd-core/wzd_site_user.h
include/wzdftpd/libwzd-core/wzd_socket.h
include/wzdftpd/libwzd-core/wzd_string.h
include/wzdftpd/libwzd-core/wzd_structs.h
include/wzdftpd/libwzd-core/wzd_threads.h
include/wzdftpd/libwzd-core/wzd_tls.h
include/wzdftpd/libwzd-core/wzd_types.h
include/wzdftpd/libwzd-core/wzd_utf8.h
include/wzdftpd/libwzd-core/wzd_vars.h
include/wzdftpd/libwzd-core/wzd_vfs.h
include/wzdftpd/libwzd.h
include/wzdftpd/libwzd_pv.h
include/wzdftpd/libwzd_socket.h
include/wzdftpd/libwzd_tls.h
lib/libwzd_core.la
lib/libwzd_core.a
lib/libwzd_core.so
lib/libwzd_core.so.1
lib/libwzd.a
lib/libwzd.la
lib/libwzd.so
lib/libwzd.so.1
share/aclocal/wzd.m4
sbin/wzdftpd
%%DATADIR%%/backends/libwzdplaintext.a
%%DATADIR%%/backends/libwzdplaintext.la
%%DATADIR%%/backends/libwzdplaintext.so
%%MYSQL%%%%DATADIR%%/backends/libwzdmysql.a
%%MYSQL%%%%DATADIR%%/backends/libwzdmysql.la
%%MYSQL%%%%DATADIR%%/backends/libwzdmysql.so
%%DATADIR%%/modules/libwzd_sfv.a
%%DATADIR%%/modules/libwzd_sfv.la
%%DATADIR%%/modules/libwzd_sfv.so
%%PERL%%%%DATADIR%%/modules/libwzd_perl.a
%%PERL%%%%DATADIR%%/modules/libwzd_perl.la
%%PERL%%%%DATADIR%%/modules/libwzd_perl.so
%%TCL%%%%DATADIR%%/modules/libwzd_tcl.a
%%TCL%%%%DATADIR%%/modules/libwzd_tcl.la
%%TCL%%%%DATADIR%%/modules/libwzd_tcl.so
%%ZEROCONF%%%%DATADIR%%/modules/libwzd_zeroconf.so
@dirrm %%DATADIR%%/modules/
@dirrm %%DATADIR%%/backends/
@dirrm %%DATADIR%%/
@dirrm include/wzdftpd/libwzd-auth
@dirrm include/wzdftpd/libwzd-base
@dirrm include/wzdftpd/libwzd-core/arch
@dirrm include/wzdftpd/libwzd-core
@dirrm include/wzdftpd
@unexec if cmp -s %D/etc/wzd.cfg %D/etc/wzd.cfg.sample; then rm -f %D/etc/wzd.cfg; fi
@unexec if cmp -s %D/etc/users %D/etc/users.sample; then rm -f %D/etc/users; fi
etc/wzd.cfg.sample
etc/users.sample
