# ftp/wzdftpd/pkg-plist
@comment $FreeBSD$
bin/siteconfig
bin/siteuptime
bin/sitewho
etc/rc.d/wzdftpd.sh
etc/users.sample
etc/wzd.cfg.sample
etc/wzd.pem
lib/libwzd.a
lib/libwzd.so
lib/libwzd.la
sbin/wzdftpd
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%%%DOCSDIR%%/Permissions.ReadMeFirst
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/TLS.ReadMeFirst
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/DOCS_XFERLOG_FORMAT
%%PORTDOCS%%%%DOCSDIR%%/KNOWN_BUGS
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%@dirrm %%DOCSDIR%%
etc/file_ginfo.txt
etc/file_help.txt
etc/file_swho.txt
etc/file_user.txt
etc/file_group.txt
etc/file_rules.txt
etc/file_users.txt
etc/file_vfs.txt
etc/file_who.txt
etc/file_groups.txt
%%DATADIR%%/backends/libwzdplaintext.a
%%DATADIR%%/backends/libwzdplaintext.so
%%DATADIR%%/backends/libwzdplaintext.la
%%DATADIR%%/modules/libwzd_sfv.la
%%DATADIR%%/modules/libwzd_sfv.a
%%DATADIR%%/modules/libwzd_sfv.so
%%TCL%%%%DATADIR%%/modules/libwzd_tcl.a
%%TCL%%%%DATADIR%%/modules/libwzd_tcl.la
%%TCL%%%%DATADIR%%/modules/libwzd_tcl.so
@dirrm %%DATADIR%%/modules/
@dirrm %%DATADIR%%/backends/
@dirrm %%DATADIR%%/
