# net/mDNSResponder/pkg-plist
bin/mDNSClientPosix
bin/mDNSIdentify
bin/mDNSNetMonitor
bin/mDNSProxyResponderPosix
bin/mDNSResponderPosix
etc/rc.d/mdns.sh
lib/libmdns.so
lib/libmdns.so.1
sbin/mdnsd
