# mail/popa3d/files/openrc-popa3d.in
#!/sbin/openrc-run

name="popa3d"
command="%%PREFIX%%/libexec/${name}"
command_args="-D"
