# mail/popa3d/files/popa3d.in
#!/bin/sh

# PROVIDE: popa3d
# REQUIRE: LOGIN cleanvar

. /etc/rc.subr

name="popa3d"
rcvar=`set_rcvar`

command="%%PREFIX%%/libexec/${name}"
command_args="-D"
: ${popa3d_enable="NO"}

load_rc_config $name
run_rc_command "$1"

