# dns/noip/pkg-plist
bin/noip2
etc/rc.d/noip.sh
%%PORTDOCS%%share/doc/noip/README.FIRST
%%PORTDOCS%%share/doc/noip/COPYING
%%PORTDOCS%%@dirrm share/doc/noip
