# dns/noip/pkg-plist
bin/noip2
etc/rc.d/noip.sh
%%PORTDOCS%%%%DOCSDIR%%/README.FIRST
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%@dirrm %%DOCSDIR%%
