# www/ziproxy/files/ziproxy.in
#!/bin/sh

# PROVIDE: ziproxy
# REQUIRE: DAEMON
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# ziproxy_enable (bool):   Set to NO by default.
#               Set it to YES to enable ziproxy.
# ziproxy_config (path):   Set to %%PREFIX%%/etc/ziproxy/ziproxy.conf
#               by default.
#

. /etc/rc.subr

name="ziproxy"
rcvar=${name}_enable

command=%%PREFIX%%/bin/${name}
pidfile=/var/run/${name}.pid

start_precmd="start_precmd"
stop_postcmd="stop_postcmd"

load_rc_config $name

: ${ziproxy_enable="NO"}
: ${ziproxy_config="%%PREFIX%%/etc/ziproxy/ziproxy.conf"}

command_args="-d -c ${ziproxy_config} --pid-file ${pidfile}"

start_precmd() {
	rm -f ${pidfile}
}

stop_postcmd() {
	rm -f ${pidfile}
}

run_rc_command "$1"
