# sysutils/logstash/files/logstash.in
#!/bin/sh

# $FreeBSD:
# PROVIDE: logstash
# REQUIRE: LOGIN
# KEYWORD: shutdown

#
# Configuration settings for logstash in /etc/rc.conf:
#
# logstash_enable (bool):
#   Set to "NO" by default.
#   Set it to "YES" to enable logstash
#
# logstash_flags (string):
#   flags to logstash.
#
# logstash_java_flags (string):
#   flags to Java VM. default is "-Des.path.data=${logstash_elastic_datadir}"
#
# logstash_mode :
#   run logstash with pre-defined flags if logstash_flags is not set. Set to ""
#   by default, which does nothing.
#   Valid options:
#     "standalone": agent, web & elasticsearch
#     "web": Starts logstash as a web ui
#     "agent": Justs works as a log shipper 
#

. /etc/rc.subr

name="logstash"
rcvar=logstash_enable

piddir="/var/run/${name}"
pidfile="${piddir}/${name}.pid"

if [ -d $piddir ]; then
	mkdir -p $piddir
fi

logdir="/var/log"
command="/usr/sbin/daemon"

load_rc_config "${name}"

: ${logstash_enable="NO"}
: ${logstash_home="%%LOGSTASH_HOME%%"}
: ${logstash_config="%%PREFIX%%/etc/${name}/${name}.conf"}
: ${logstash_jar="%%LOGSTASH_HOME%%/%%LOGSTASH_JAR%%"}
: ${logstash_java_home="%%JAVA_HOME%%"}
: ${logstash_log="NO"}
: ${logstash_mode=""}
: ${logstash_port="9292"}
: ${logstash_elastic_backend=""}
: ${logstash_log_file="${logdir}/${name}.log"}
: ${logstash_elastic_datadir="%%LOGSTASH_DATA_DIR%%"}
: ${logstash_log_options="--log ${logstash_log_file}"}
: ${logstash_java_flags="-Des.path.data=${logstash_elastic_datadir}"}

java_cmd="${logstash_java_home}/bin/java"
procname="${java_cmd}"

logstash_chdir=${logstash_home}

if [ -z ${logstash_flags} ]; then
    if [ ${logstash_mode} == "standalone" ]; then
        logstash_flags="agent -f ${logstash_config} -- web --port ${logstash_port} --backend elasticsearch:///?local ${logstash_log_options}"
    elif [ ${logstash_mode} == "agent" ];then
        logstash_flags="agent -f ${logstash_config} ${logstash_log_options}"
    elif [ ${logstash_mode} == "web" ];then
        logstash_flags="web --port ${logstash_port} --backend elasticsearch://${logstash_elastic_backend}/ ${logstash_log_options}"
    else
        logstash_flags="-f ${logstash_config}"
    fi
fi

command_args="-f -p ${pidfile} ${java_cmd} ${logstash_java_flags} -jar ${logstash_jar} ${logstash_flags}"
required_files="${java_cmd} ${logstash_config}"

run_rc_command "$1"
