# security/xinetd/pkg-plist
bin/xconv
etc/rc.d/xinetd.sh
sbin/xinetd
