# security/xinetd/pkg-plist
etc/rc.d/xinetd.sh
man/man1/xinetd.1.gz
man/man5/xinetd.conf.5.gz
man/man5/xinetd.log.5.gz
sbin/xinetd
