# net/dictd/files/openrc-dictd.in
#!/sbin/openrc-run

name=dictd
command=%%PREFIX%%/sbin/dictd
