# multimedia/ffmpeg-devel/files/ffserver-head.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ffserver-head
# REQUIRE: NETWORKING
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable ffserver:
#
#ffserver_head_enable="YES"
#
ffserver_head_enable="${ffserver_head_enable-NO}"

. /etc/rc.subr

name=ffserver_head
rcvar=ffserver_head_enable

command="%%PREFIX%%/bin/${name} &"
required_files=%%PREFIX%%/etc/ffserver-head.conf

load_rc_config ${name}
run_rc_command "$1"
