# sysutils/runit/pkg-plist
@comment $FreeBSD$
etc/rc.d/runsvdir.sh.sample
etc/runit/1.sample
etc/runit/2.sample
etc/runit/3.sample
etc/runit/ctrlaltdel.sample
etc/runit/getty-ttyv4/run.sample
etc/runit/getty-ttyv4/finish.sample
sbin/chpst
sbin/runit
sbin/runit-init
sbin/runsv
sbin/runsvchdir
sbin/runsvctrl
sbin/runsvdir
sbin/runsvstat
sbin/svlogd
sbin/svwaitdown
sbin/svwaitup
sbin/utmpset
%%PORTDOCS%%%%DOCSDIR%%/CHANGES
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/THANKS
%%PORTDOCS%%%%DOCSDIR%%/benefits.html
%%PORTDOCS%%%%DOCSDIR%%/chpst.8.html
%%PORTDOCS%%%%DOCSDIR%%/dependencies.html
%%PORTDOCS%%%%DOCSDIR%%/index.html
%%PORTDOCS%%%%DOCSDIR%%/install.html
%%PORTDOCS%%%%DOCSDIR%%/replaceinit.html
%%PORTDOCS%%%%DOCSDIR%%/runit-init.8.html
%%PORTDOCS%%%%DOCSDIR%%/runit.8.html
%%PORTDOCS%%%%DOCSDIR%%/runlevels.html
%%PORTDOCS%%%%DOCSDIR%%/runscripts.html
%%PORTDOCS%%%%DOCSDIR%%/runsv.8.html
%%PORTDOCS%%%%DOCSDIR%%/runsvchdir.8.html
%%PORTDOCS%%%%DOCSDIR%%/runsvctrl.8.html
%%PORTDOCS%%%%DOCSDIR%%/runsvdir.8.html
%%PORTDOCS%%%%DOCSDIR%%/runsvstat.8.html
%%PORTDOCS%%%%DOCSDIR%%/svlogd.8.html
%%PORTDOCS%%%%DOCSDIR%%/svwaitdown.8.html
%%PORTDOCS%%%%DOCSDIR%%/svwaitup.8.html
%%PORTDOCS%%%%DOCSDIR%%/upgrade.html
%%PORTDOCS%%%%DOCSDIR%%/usedietlibc.html
%%PORTDOCS%%%%DOCSDIR%%/useinit.html
%%PORTDOCS%%%%DOCSDIR%%/utmpset.8.html
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm etc/runit/getty-ttyv4
@dirrm etc/runit
