# devel/dbus/files/dbus.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: dbus
# REQUIRE: DAEMON ldconfig
#
# Add the following lines to /etc/rc.conf to enable the D-BUS messaging system:
#
# dbus_enable="YES"
#

. /etc/rc.subr

: ${dbus_enable=${gnome_enable-NO}} ${dbus_flags="--system"}

name=dbus
rcvar=dbus_enable

command="%%PREFIX%%/bin/dbus-daemon"
pidfile="/var/run/dbus/pid"

start_precmd="dbus_prestart"
stop_postcmd="dbus_poststop"

dbus_prestart()
{
    %%PREFIX%%/bin/dbus-uuidgen --ensure
}

dbus_poststop()
{
    rm -f $pidfile
# The following two lines may be removed after 2018-01-01
    rm -f /var/db/dbus/machine-id
    [ ! -d /var/db/dbus ] || rmdir /var/db/dbus
}

load_rc_config ${name}
run_rc_command "$1"
