# devel/dbus/files/openrc-dbus.in
#!/sbin/openrc-run
# Copyright (c) 2007-2015 The OpenRC Authors.
# See the Authors file at the top-level directory of this distribution and
# https://github.com/OpenRC/openrc/blob/master/AUTHORS
#
# This file is part of OpenRC. It is subject to the license terms in
# the LICENSE file found in the top-level directory of this
# distribution and at https://github.com/OpenRC/openrc/blob/master/LICENSE
# This file may not be copied, modified, propagated, or distributed
# except according to the terms contained in the LICENSE file.

command=%%PREFIX%%/bin/dbus-daemon
pidfile=/var/run/dbus.pid
command_args="${dbusd_args---system} --nofork"
name="dbus"
supervisor="supervise-daemon"
supervise_daemon_args="-1 /var/log/dbus.out -2 /var/log/dbus.err"

system_bus_socket="/var/run/dbus/system_bus_socket"
system_bus_pid="/var/run/dbus/pid"

depend()
{
	need localmount
	after bootmisc
}

start_pre()
{
	mkdir -p $(dirname $pidfile)
	if [ ! -e "/var/db/dbus/machine-id" ] ; then
	/usr/local/bin/dbus-uuidgen --ensure
	fi
}

start_post()
{
	#it takes a moment for the system bus to get initialized
	# pause until it is actually ready before returning
	# maximum of 5 seconds wait
	for i in 1 2 3 4 5
	do
		if [ -S ${system_bus_socket} ] ; then
		  #einfo "System bus ready"
		  break
		else
		  #einfo "Waiting for system bus to initialize..."
		  sleep 1
		fi
	done
	return 0
}

stop_post()
{
  if [ -S ${system_bus_socket} ] ; then
    rm ${system_bus_socket}
    rm ${system_bus_pid}
  fi
  return 0
}
