# www/calendarserver/files/caldavd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: caldavd
# REQUIRE: LOGIN
#
# Add the following lines to /etc/rc.conf to run caldavd:
#
# caldavd_enable (bool):           Set it to "YES" to enable caldavd.
#                                Default is "NO".
# caldavd_conf (file):             Local config file.
#                                Default is "/usr/local/etc/caldavd.conf".
# caldavd_flags (string):          Flags to pass to caldavd.
#                                Default is empty.
# caldavd_pidfile (string):        Set path to the pid file.
#                                Default is /var/run/caldavd/caldavd.pid
#

. /etc/rc.subr

name="caldavd"
rcvar=caldavd_enable

load_rc_config $name
: ${caldavd_enable:=no}
: ${caldavd_conf="%%PREFIX%%/etc/caldavd/caldavd.plist"}
: ${caldavd_pidfile="/var/run/caldavd/${name}.pid"}

pidfile="${caldavd_pidfile}"
required_files="${caldavd_conf}"

procname="%%PYTHON_CMD%%"
PATH=${PATH}:%%PREFIX%%/bin

command="%%PREFIX%%/bin/caldavd"
command_args=" -f ${caldavd_conf} -R kqueue"

run_rc_command "$1"
