# www/calendarserver/files/caldavd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: caldavd
# REQUIRE: LOGIN
#
# Add the following lines to /etc/rc.conf to run caldavd:
#
# caldavd_enable (bool):           Set it to "YES" to enable caldavd.
#                                Default is "NO".
# caldavd_conf (file):             Local config file.
#                                Default is "/usr/local/etc/caldavd.conf".
# caldavd_flags (string):          Flags to pass to caldavd.
#                                Default is empty.
#

. /etc/rc.subr

: ${caldavd_enable="NO"}
: ${caldavd_conf="%%PREFIX%%/etc/caldavd/caldavd.plist"}


name="caldavd"

command="%%PREFIX%%/bin/caldavd"
command_args=" -f ${caldavd_conf} -R kqueue"

pidfile="/var/db/${name}/${name}.pid"
required_files="${caldavd_conf}"

procname="%%PYTHON_CMD%%"    
PATH=${PATH}:%%PREFIX%%/bin

start_precmd="install -d -o caldavd /var/db/caldavd/Data/Documents /var/db/caldavd/Documents && install -d /var/log/caldavd"

run_rc_command "$1"
