# net/pipsecd/pkg-plist
sbin/pipsecd
etc/ipsec/pipsecd.conf.sample
etc/rc.d/pipsecd.sh
%%PORTDOCS%%share/doc/pipsecd/README
%%PORTDOCS%%share/doc/pipsecd/COPYRIGHT
%%PORTDOCS%%share/doc/pipsecd/README.port
%%PORTDOCS%%@dirrm share/doc/pipsecd
@dirrm etc/ipsec
