# net/pipsecd/pkg-plist
sbin/pipsecd
etc/ipsec/pipsecd.conf.sample
etc/rc.d/pipsecd.sh
share/doc/pipsecd/README
share/doc/pipsecd/COPYRIGHT
share/doc/pipsecd/README.port
@dirrm share/doc/pipsecd
@dirrm etc/ipsec
