# devel/dcons/pkg-plist
sbin/dconschat
@unexec [ -f %%MAKEDEV%% ] && /bin/rm -f %%DCONS_DEV%% || true
@unexec /bin/rm -f %D/lib/dcons/linker.hints
lib/dcons/dcons.ko
lib/dcons/dcons_crom.ko
@exec type kldxref > /dev/null 2>&1 && kldxref %D/lib/dcons
@exec [ -f %%MAKEDEV%% ] && /bin/rm -f %%DCONS_DEV%%;    %%MKNOD%% %%DCONS_DEV%%    c %%DCONS_MAJOR%%    %%DCONS_MINOR%%;    /bin/chmod %%DCONS_MODE%%    || true
etc/rc.d/dcons.sh
%%DOCSDIR%%/gdb-usage
%%DOCSDIR%%/README.dcons
@dirrm lib/dcons
@dirrm %%DOCSDIR%%
