# net/dictd/pkg-plist
bin/dictzip
sbin/dictd
share/dictd/example.conf
share/dictd/example.dictrc
share/dictd/example.site
share/dictd/example2.conf
share/dictd/example3.conf
etc/dictd.conf
etc/rc.d/dictd.sh
@dirrm share/dictd
