# sysutils/salt/files/salt_minion.in
#!/bin/sh

# $FreeBSD$
#
# Salt Minion startup script
#
# PROVIDE: salt_minion
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following to /etc/rc.conf[.local] to enable this service
#
# salt_minion_enable="YES"
#

. /etc/rc.subr

name="salt_minion"
rcvar="${name}_enable"

load_rc_config ${name}

: ${salt_minion_enable:="NO"}

command="%%PREFIX%%/bin/salt-minion"
command_interpreter="%%PYTHON_CMD%%"
required_files="%%PREFIX%%/etc/salt/minion"
command_args="-c ${required_files} -d"

run_rc_command "$1"
