# comms/thebridge/pkg-plist
bin/tbdcmd
etc/rc.d/tbd.sh.sample
etc/tbd.conf.sample
libexec/tbd
