# www/wwwoffle/files/openrc-wwwoffled.in
#!/sbin/openrc-run

# Define these wwwoffled_* variables in one of these files:
#	/etc/rc.conf
#	/etc/rc.conf.local
#	/etc/rc.conf.d/wwwoffled
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#
wwwoffled_flags=${wwwoffled_flags-""}
wwwoffled_pidfile=${wwwoffled_pidfile-"/var/run/wwwoffled.pid"}
wwwoffled_conffile="%%PREFIX%%/etc/wwwoffle.conf"

name="wwwoffled"
command="%%PREFIX%%/sbin/wwwoffled"
pidfile="${wwwoffled_pidfile}"
required_files="${wwwoffled_conffile}"

depend(){
	keyword -shutdown
}
