# mail/masqmail/pkg-plist
etc/rc.d/masqmail.sh
etc/masqmail/masqmail.conf-dist
etc/masqmail/masqmail.get-dist
etc/masqmail/masqmail.route-dist
sbin/masqmail
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/THOUGHTS
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%%%DOCSDIR%%/install.html
@dirrm etc/masqmail
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@unexec rmdir /var/log/masqmail 2>/dev/null || true
