# japanese/wxg/pkg-plist
bin/wxgdic
bin/wxgdmp
bin/wxgdset
etc/rc.d/wxg.sh.sample
lib/wxg/dic/98dic.txt
lib/wxg/dic/98txdic.txt
lib/wxg/dic/chimei.gdj
lib/wxg/dic/english.gdj
lib/wxg/dic/kao.gdj
lib/wxg/dic/seimei.gdj
lib/wxg/dic/tel.gdj
lib/wxg/dic/wxg.gdf
lib/wxg/dic/wxg.gdj
lib/wxg/dic/wxg.gds
lib/wxg/dic/wxgrev.gdf
lib/wxg/dic/wxgsei.gds
lib/wxg/dic/wxgseiu.gds
lib/wxg/dic/wxgsp.gdf
lib/wxg/dic/wxgtx.gdj
lib/wxg/dic/wxgu.gdj
lib/wxg/dic/wxgu.gds
lib/wxg/dic/wxgurev.gdj
lib/wxg/dic/yubin.gdj
lib/wxg/roma/wxgroma.grt
lib/wxg/wxg.gct
lib/wxg/wxgrc
@dirrm lib/wxg/roma
@dirrm lib/wxg/dic
@dirrm lib/wxg
sbin/wxgserver
%%PORTDOCS%%share/doc/wxg/README
%%PORTDOCS%%share/doc/wxg/html/ailogo2.gif
%%PORTDOCS%%share/doc/wxg/html/appendix.html
%%PORTDOCS%%share/doc/wxg/html/contract.html
%%PORTDOCS%%share/doc/wxg/html/custom.html
%%PORTDOCS%%share/doc/wxg/html/fcontents.html
%%PORTDOCS%%share/doc/wxg/html/hinshi.html
%%PORTDOCS%%share/doc/wxg/html/imi.html
%%PORTDOCS%%share/doc/wxg/html/index.html
%%PORTDOCS%%share/doc/wxg/html/inst.html
%%PORTDOCS%%share/doc/wxg/html/intro.html
%%PORTDOCS%%share/doc/wxg/html/jinput.html
%%PORTDOCS%%share/doc/wxg/html/relcanna.html
%%PORTDOCS%%share/doc/wxg/html/uninst.html
%%PORTDOCS%%share/doc/wxg/html/uty.html
%%PORTDOCS%%@dirrm share/doc/wxg/html
%%PORTDOCS%%@dirrm share/doc/wxg
@unexec [ -L /usr/local/lib/wxg ] && rm /usr/local/lib/wxg || true
@exec [ X"%D" != X"/usr/local" ] && ln -s %D/lib/wxg /usr/local/lib/ || true
