# multimedia/vdr/files/vdr.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: vdr
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf[.local] to enable vdr:
#
# vdr_enable="YES"
#
# And if you are using vdr plugins etc. add those including any plugin
# parameters to /etc/rc.conf[.local] similar to this:
#
# vdr_flags="'-Pxineliboutput --local=none --remote=127.1:37890' \
#	-Pfemon -Posdpip \
#	-Posdteletext --cache-system=packed --max-cache=128' \
#	-Pstreamdev-server '-Plive -i127.1' -Pepgsearch \
#	-Pinfosatepg"

. /etc/rc.subr

name=vdr
rcvar=vdr_enable

load_rc_config ${name}
: ${vdr_enable="NO"}
: ${vdr_chdir="%%VIDEODIR%%"}
: ${vdr_waitdvb="YES"}

command="%%PREFIX%%/bin/vdr"
procname="$command"
command_args="-d"

start_precmd="${name}_prestart"

vdr_prestart()
{
	if checkyesno "${name}_waitdvb"; then
		command="/usr/sbin/daemon"
		rc_flags="-f %%PREFIX%%/bin/vdr-waitdvb $rc_flags"
	fi
}

run_rc_command "$1"
