# net-mgmt/bigsister/files/bigsister.sh.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: bigsister
# REQUIRE: DAEMON
# BEFORE:  LOGIN
#
# Add the following to /etc/rc.conf to enable this service
#
# bigsister_enable="YES"
#

. /etc/rc.subr

name="bigsister"
rcvar=bigsister_enable
start_cmd=bs_start
stop_cmd=bs_stop

bs_start() {
	%%PREFIX%%/bigsister/bin/bb_start start
}

bs_stop() {
	%%PREFIX%%/bigsister/bin/bb_start stop
}

load_rc_config $name
run_rc_command "$1"
