# dns/maradns/files/Deadwood.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: Deadwood
# REQUIRE: SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown
#

. /etc/rc.subr

: ${dwood3_enable="NO"}
: ${dwood3_conf="%%ETCDIR%%/dwood3rc"}

name="Deadwood"
rcvar=dwood3_enable
command="%%PREFIX%%/bin/duende"
pidfile="/var/run/${name}.pid"
command_args="--pid=${pidfile} %%PREFIX%%/sbin/Deadwood -f ${dwood3_conf}"

load_rc_config $name

run_rc_command "$1"
