# dns/maradns/files/maradns.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: maradns
# REQUIRE: SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown
#

. /etc/rc.subr

: ${maradns_enable="NO"}
: ${maradns_conf="%%ETCDIR%%/mararc"}

name="maradns"
rcvar=maradns_enable
command="%%PREFIX%%/bin/duende"
pidfile="/var/run/${name}.pid"
command_args="--pid=${pidfile} %%PREFIX%%/sbin/maradns -f ${maradns_conf}"

load_rc_config $name

run_rc_command "$1"
