# dns/maradns/files/openrc-maradns.in
#!/sbin/openrc-run

: ${maradns_conf="%%ETCDIR%%/mararc"}

name="maradns"
command="%%PREFIX%%/bin/duende"
pidfile="/var/run/${name}.pid"
command_args="--pid=${pidfile} %%PREFIX%%/sbin/maradns -f ${maradns_conf}"

depend() {
	keyword -shutdown
}
