# dns/maradns/files/zoneserver.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: zoneserver
# REQUIRE: SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown
#

. /etc/rc.subr

: ${zoneserver_enable="NO"}
: ${maradns_conf="%%PREFIX%%/etc/mararc"}

name="zoneserver"
rcvar=`set_rcvar`
command="%%PREFIX%%/bin/duende"
command_args="%%PREFIX%%/sbin/zoneserver -f ${maradns_conf}"
pidfile="/var/run/${name}.pid"

load_rc_config $name

run_rc_command "$1"
