# dns/maradns/files/zoneserver.in
#!/bin/sh

# PROVIDE: zoneserver
# REQUIRE: SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown
#

. /etc/rc.subr

: ${zoneserver_enable="NO"}
: ${maradns_conf="%%ETCDIR%%/mararc"}

name="zoneserver"
rcvar=zoneserver_enable
command="%%PREFIX%%/bin/duende"
pidfile="/var/run/${name}.pid"
command_args="--pid=${pidfile} %%PREFIX%%/sbin/zoneserver -f ${maradns_conf}"

load_rc_config $name

run_rc_command "$1"
