# comms/comserv/pkg-plist
sbin/comservd
etc/comservd.conf.sample
etc/rc.d/comservd.sh.sample
