# net/isc-dhcp3/pkg-plist
@comment $FreeBSD$
bin/omshell
etc/dhclient.conf.sample
@exec [ -f %D/etc/dhclient.conf ] || touch %D/etc/dhclient.conf
@unexec [ -L %D/etc/dhclient.conf -o -s %D/etc/dhclient.conf ] || rm -f %D/etc/dhclient.conf
etc/dhcpd.conf.sample
@unexec if cmp -s %D/etc/rc.isc-dhcpd.conf %D/etc/rc.isc-dhcpd.conf.sample; then rm -f %D/etc/rc.isc-dhcpd.conf; fi
etc/rc.isc-dhcpd.conf.sample
@exec [ -f %D/etc/rc.isc-dhcpd.conf ] || cp %D/etc/rc.isc-dhcpd.conf.sample %D/etc/rc.isc-dhcpd.conf
@unexec if cmp -s %D/etc/rc.isc-dhcrelay.conf %D/etc/rc.isc-dhcrelay.conf.sample; then rm -f %D/etc/rc.isc-dhcrelay.conf; fi
etc/rc.isc-dhcrelay.conf.sample
@exec [ -f %D/etc/rc.isc-dhcrelay.conf ] || cp %D/etc/rc.isc-dhcrelay.conf.sample %D/etc/rc.isc-dhcrelay.conf
etc/rc.d/isc-dhcpd.sh.sample
@unexec rm -f etc/rc.d/isc-dhcpd.sh
etc/rc.d/isc-dhcrelay.sh.sample
@unexec rm -f etc/rc.d/isc-dhcrelay.sh
include/dhcpctl.h
include/isc-dhcp/boolean.h
include/isc-dhcp/dst.h
include/isc-dhcp/int.h
include/isc-dhcp/lang.h
include/isc-dhcp/list.h
include/isc-dhcp/result.h
include/isc-dhcp/types.h
include/omapip/alloc.h
include/omapip/buffer.h
include/omapip/omapip.h
lib/libdhcpctl.a
lib/libomapi.a
sbin/dhclient
sbin/dhclient-script
sbin/dhcpd
sbin/dhcrelay
%%PORTDOCS%%share/doc/isc-dhcp3/ANONCVS
%%PORTDOCS%%share/doc/isc-dhcp3/CHANGES
%%PORTDOCS%%share/doc/isc-dhcp3/COPYRIGHT
%%PORTDOCS%%share/doc/isc-dhcp3/README
%%PORTDOCS%%share/doc/isc-dhcp3/RELNOTES
%%PORTDOCS%%@dirrm share/doc/isc-dhcp3
@dirrm include/omapip
@dirrm include/isc-dhcp
@exec [ -f /var/run/dhclient.leases ] || touch /var/run/dhclient.leases
@unexec [ -s /var/run/dhclient.leases ] || rm -f /var/run/dhclient.leases
@exec [ -f /var/run/dhcpd.leases ] || touch /var/run/dhcpd.leases
@unexec [ -s /var/run/dhcpd.leases ] || rm -f /var/run/dhcpd.leases
