# sysutils/usermin/files/openrc-usermin.in
#!/sbin/openrc-run

name=usermin

procname=%%PERL%%
pidfile=/var/log/usermin/miniserv.pid
required_dirs=%%PREFIX%%/etc/usermin
command=%%PREFIX%%/etc/usermin/start

depend() {
	keyword -shutdown
}
