# sysutils/hal/pkg-plist
@unexec %D/etc/rc.d/hald stop > /dev/null 2>&1 || true
bin/hal-device
bin/hal-disable-polling
bin/hal-find-by-capability
bin/hal-find-by-property
bin/hal-get-property
bin/hal-is-caller-locked-out
bin/hal-is-caller-privileged
bin/hal-lock
bin/hal-set-property
bin/lshal
@exec mkdir -p %D/etc/hal/fdi/preprobe
@exec mkdir -p %D/etc/hal/fdi/policy
@exec mkdir -p %D/etc/hal/fdi/information
include/hal/libhal-storage.h
include/hal/libhal.h
lib/libhal-storage.a
lib/libhal-storage.la
lib/libhal-storage.so
lib/libhal-storage.so.1
lib/libhal.a
lib/libhal.la
lib/libhal.so
lib/libhal.so.1
libdata/pkgconfig/hal-storage.pc
libdata/pkgconfig/hal.pc
libexec/hal-storage-cleanup-all-mountpoints
libexec/hal-storage-cleanup-mountpoint
libexec/hal-storage-closetray
libexec/hal-storage-eject
libexec/hal-storage-mount
libexec/hal-storage-unmount
libexec/hal-system-power-pm-is-supported
libexec/hal-system-setserial
libexec/hal/scripts/freebsd/hal-system-lcd-get-brightness-freebsd
libexec/hal/scripts/freebsd/hal-system-lcd-set-brightness-freebsd
libexec/hal/scripts/freebsd/hal-system-power-reboot-freebsd
libexec/hal/scripts/freebsd/hal-system-power-set-power-save-freebsd
libexec/hal/scripts/freebsd/hal-system-power-shutdown-freebsd
libexec/hal/scripts/freebsd/hal-system-power-suspend-freebsd
libexec/hal/scripts/hal-dockstation-undock
libexec/hal/scripts/hal-functions
libexec/hal/scripts/hal-luks-remove
libexec/hal/scripts/hal-luks-setup
libexec/hal/scripts/hal-luks-teardown
libexec/hal/scripts/hal-system-killswitch-get-power
libexec/hal/scripts/hal-system-killswitch-set-power
libexec/hal/scripts/hal-system-lcd-get-brightness
libexec/hal/scripts/hal-system-lcd-set-brightness
libexec/hal/scripts/hal-system-power-hibernate
libexec/hal/scripts/hal-system-power-reboot
libexec/hal/scripts/hal-system-power-set-power-save
libexec/hal/scripts/hal-system-power-shutdown
libexec/hal/scripts/hal-system-power-suspend
libexec/hal/scripts/hal-system-power-suspend-hybrid
libexec/hal/scripts/hal-system-wol-enable
libexec/hal/scripts/hal-system-wol-enabled
libexec/hal/scripts/hal-system-wol-supported
libexec/hald-addon-mouse-sysmouse
libexec/hald-addon-storage
libexec/hald-generate-fdi-cache
libexec/hald-probe-hiddev
libexec/hald-probe-mouse
libexec/hald-probe-scsi
libexec/hald-probe-smbios
libexec/hald-probe-storage
%%USB2%%libexec/hald-probe-usb2-device
%%USB2%%libexec/hald-probe-usb2-interface
%%USB2%%libexec/hald-probe-video4linux
libexec/hald-probe-volume
libexec/hald-runner
sbin/hald
%%PORTDOCS%%%%DOCSDIR%%/README.freebsd
%%PORTDOCS%%%%DOCSDIR%%/README.fuse
%%DATADIR%%/dist/hal.conf
%%DATADIR%%/fdi/fdi.dtd
%%DATADIR%%/fdi/information/10freedesktop/01-deprecated-keys.fdi
@exec mkdir -p %D/%%DATADIR%%/fdi/information/20thirdparty
%%DATADIR%%/fdi/policy/10osvendor/10-dockstation.fdi
%%DATADIR%%/fdi/policy/10osvendor/10-input-policy.fdi
%%DATADIR%%/fdi/policy/10osvendor/10-laptop-panel-mgmt-policy.fdi
%%DATADIR%%/fdi/policy/10osvendor/10-leds.fdi
%%DATADIR%%/fdi/policy/10osvendor/10-mouse-sysmouse.fdi
%%DATADIR%%/fdi/policy/10osvendor/10-power-mgmt-policy.fdi
%%DATADIR%%/fdi/policy/10osvendor/10-rfkill-switch.fdi
%%DATADIR%%/fdi/policy/10osvendor/10-tabletPCs.fdi
%%DATADIR%%/fdi/policy/10osvendor/10-x11-input.fdi
%%DATADIR%%/fdi/policy/10osvendor/15-storage-luks.fdi
%%DATADIR%%/fdi/policy/10osvendor/20-storage-methods.fdi
%%DATADIR%%/fdi/policy/10osvendor/30-wol.fdi
%%DATADIR%%/mount-fuse
share/PolicyKit/policy/org.freedesktop.hal.dockstation.policy
share/PolicyKit/policy/org.freedesktop.hal.killswitch.policy
share/PolicyKit/policy/org.freedesktop.hal.leds.policy
share/PolicyKit/policy/org.freedesktop.hal.policy
share/PolicyKit/policy/org.freedesktop.hal.power-management.policy
share/PolicyKit/policy/org.freedesktop.hal.storage.policy
share/PolicyKit/policy/org.freedesktop.hal.wol.policy
@exec mkdir -p %D/%%DATADIR%%/fdi/policy/20thirdparty
@exec mkdir -p %D/%%DATADIR%%/fdi/preprobe/10osvendor
@exec mkdir -p %D/%%DATADIR%%/fdi/preprobe/20thirdparty
@dirrm %%DATADIR%%/fdi/preprobe/20thirdparty
@dirrm %%DATADIR%%/fdi/preprobe/10osvendor
@dirrm %%DATADIR%%/fdi/preprobe
@dirrm %%DATADIR%%/fdi/policy/20thirdparty
@dirrm %%DATADIR%%/fdi/policy/10osvendor
@dirrm %%DATADIR%%/fdi/policy
@dirrm %%DATADIR%%/fdi/information/20thirdparty
@dirrm %%DATADIR%%/fdi/information/10freedesktop
@dirrm %%DATADIR%%/fdi/information
@dirrm %%DATADIR%%/fdi
@dirrm %%DATADIR%%/dist
@dirrm %%DATADIR%%
@dirrm %%DOCSDIR%%
@dirrm libexec/hal/scripts/freebsd
@dirrm libexec/hal/scripts
@dirrm libexec/hal
@dirrm include/hal
@dirrm etc/hal/fdi/preprobe
@dirrm etc/hal/fdi/policy
@dirrm etc/hal/fdi/information
@dirrm etc/hal/fdi
@dirrm etc/hal
@unexec rm -f /var/run/hald/hald.pid 2>/dev/null || true
@unexec rm -rf /var/run/hald/hald-local 2>/dev/null || true
@unexec rm -rf /var/run/hald/hald-runner 2>/dev/null || true
@unexec rmdir /var/run/hald 2>/dev/null || true
@unexec rmdir /var/cache/hald 2>/dev/null || true
@unexec rmdir /var/cache 2>/dev/null || true
@unexec rmdir /var/lib/hal 2>/dev/null || true
@unexec rmdir /var/lib 2>/dev/null || true
