# print/pips800/pkg-plist
%%EKPNAVI%%bin/ekpnavi
bin/ekpstm
%%GSCONFIG%%bin/gsconfig%%PRTYPE%%
bin/pips%%PRTYPE%%
%%CUPS%%libexec/cups/backend/ekplp
%%CUPS%%libexec/cups/filter/pipstoprinter
%%CUPS%%libexec/cups/filter/rastertopips
%%CUPS%%libexec/pips/%%PRT_MODEL%%/cupsopt.csv
libexec/pips/%%PRT_MODEL%%/filter%%PRTYPE%%
libexec/pips/%%PRT_MODEL%%/lang
libexec/pips/%%PRT_MODEL%%/setup
libexec/pips/%%PRT_MODEL%%/scripts/en.lc
libexec/pips/%%PRT_MODEL%%/scripts/ja.lc
@dirrm libexec/pips/%%PRT_MODEL%%/scripts
%%DTRFILTER%%libexec/pips/printer/dtrfilter
libexec/pips/printer/ekpd
%%DTRFILTER%%libexec/pips/printer/freset
%%DTRFILTER%%libexec/pips/printer/gsconfig
%%DTRFILTER%%libexec/pips/printer/paper_list.csv
@dirrm libexec/pips/printer
%%CUPS%%share/cups/model/%%MODEL_FILE%%
%%PORTDOCS%%share/doc/pips%%PRTYPE%%/COPYING
%%PORTDOCS%%share/doc/pips%%PRTYPE%%/COPYING.KOWA
%%PORTDOCS%%share/doc/pips%%PRTYPE%%/COPYING.KOWA.ja
%%PORTDOCS%%share/doc/pips%%PRTYPE%%/COPYING.LIB
share/doc/pips%%PRTYPE%%/readme%%PRTYPE%%
%%CUPS%%share/doc/pips%%PRTYPE%%/readme%%PRTYPE%%-cups
share/doc/pips%%PRTYPE%%/%%LIB_README%%
@dirrm share/doc/pips%%PRTYPE%%
etc/rc.d/ekpd.sh
etc/ekpdrc
etc/pipsrc
@exec chmod 666 %D/etc/pipsrc
share/locale/de/LC_MESSAGES/pips%%PRTYPE%%.mo
share/locale/es/LC_MESSAGES/pips%%PRTYPE%%.mo
share/locale/fr/LC_MESSAGES/pips%%PRTYPE%%.mo
share/locale/it/LC_MESSAGES/pips%%PRTYPE%%.mo
share/locale/ja/LC_MESSAGES/pips%%PRTYPE%%.mo
share/locale/ko/LC_MESSAGES/pips%%PRTYPE%%.mo
share/locale/nl/LC_MESSAGES/pips%%PRTYPE%%.mo
share/locale/pt/LC_MESSAGES/pips%%PRTYPE%%.mo
share/locale/zh/LC_MESSAGES/pips%%PRTYPE%%.mo
share/locale/zh_TW/LC_MESSAGES/pips%%PRTYPE%%.mo
%%LPR%%share/locale/de/LC_MESSAGES/ekpnavi.mo
%%LPR%%share/locale/es/LC_MESSAGES/ekpnavi.mo
%%LPR%%share/locale/fr/LC_MESSAGES/ekpnavi.mo
%%LPR%%share/locale/it/LC_MESSAGES/ekpnavi.mo
%%LPR%%share/locale/ja/LC_MESSAGES/ekpnavi.mo
%%LPR%%share/locale/ko/LC_MESSAGES/ekpnavi.mo
%%LPR%%share/locale/nl/LC_MESSAGES/ekpnavi.mo
%%LPR%%share/locale/pt/LC_MESSAGES/ekpnavi.mo
%%LPR%%share/locale/zh/LC_MESSAGES/ekpnavi.mo
%%LPR%%share/locale/zh_TW/LC_MESSAGES/ekpnavi.mo
%%LPR%%share/locale/de/LC_MESSAGES/ekpstm.mo
%%LPR%%share/locale/es/LC_MESSAGES/ekpstm.mo
%%LPR%%share/locale/fr/LC_MESSAGES/ekpstm.mo
%%LPR%%share/locale/it/LC_MESSAGES/ekpstm.mo
%%LPR%%share/locale/ja/LC_MESSAGES/ekpstm.mo
%%LPR%%share/locale/ko/LC_MESSAGES/ekpstm.mo
%%LPR%%share/locale/nl/LC_MESSAGES/ekpstm.mo
%%LPR%%share/locale/pt/LC_MESSAGES/ekpstm.mo
%%LPR%%share/locale/zh/LC_MESSAGES/ekpstm.mo
%%LPR%%share/locale/zh_TW/LC_MESSAGES/ekpstm.mo
