# print/pips800/pkg-plist
%%EKPNAVI%%bin/ekpnavi
bin/ekpstm
%%GSCONFIG%%bin/gsconfig%%PRTYPE%%
%%PIPS%%bin/pips%%PRTYPE%%
%%CUPS%%libexec/cups/backend/ekplp
%%CUPS%%libexec/cups/filter/pipstoprinter
%%CUPS%%libexec/cups/filter/rastertopips
%%CUPS%%libexec/pips/%%PRT_MODEL%%/%%CUPSOPT_FILE%%
%%PIPS%%libexec/pips/%%PRT_MODEL%%/filter%%PRTYPE%%
%%PIPS%%libexec/pips/%%PRT_MODEL%%/lang
libexec/pips/%%PRT_MODEL%%/setup
libexec/pips/%%PRT_MODEL%%/scripts/en.lc
libexec/pips/%%PRT_MODEL%%/scripts/ja.lc
@dirrm libexec/pips/%%PRT_MODEL%%/scripts
%%DTRFILTER%%libexec/pips/printer/dtrfilter
libexec/pips/printer/ekpd
%%DTRFILTER%%libexec/pips/printer/freset
%%DTRFILTER%%libexec/pips/printer/gsconfig
%%PAPER_LIST%%libexec/pips/printer/paper_list.csv
@dirrm libexec/pips/printer
%%CUPS%%share/cups/model/%%MODEL_FILE%%
%%PORTDOCS%%share/doc/pips%%PRTYPE%%/COPYING
%%PORTDOCS%%share/doc/pips%%PRTYPE%%/COPYING.KOWA
%%PORTDOCS%%share/doc/pips%%PRTYPE%%/COPYING.KOWA.ja
%%PORTDOCS%%share/doc/pips%%PRTYPE%%/COPYING.LIB
%%PIPS%%share/doc/pips%%PRTYPE%%/readme%%PRTYPE%%
%%CUPS%%share/doc/pips%%PRTYPE%%/readme%%PRTYPE%%-cups
share/doc/pips%%PRTYPE%%/%%LIB_README%%
@dirrm share/doc/pips%%PRTYPE%%
etc/rc.d/ekpd.sh
etc/ekpdrc
%%PIPS%%etc/pipsrc
%%PIPS%%@exec chmod 666 %D/etc/pipsrc
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/Makefile
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/Makefile.am
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/Makefile.in
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/c.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/dy.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/error.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/error_i3.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/error_i4.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/error_i6.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/error_i8.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/error_ia.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/k.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/lc.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/lm.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/m.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/printing1.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/printing2.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/printing3.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/ready.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/ready_i1.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/warning.xpm
%%EKPSTM_PIXMAP%%share/ekpstm/pixmaps/y.xpm
%%EKPSTM_PIXMAP%%@dirrm share/ekpstm/pixmaps
%%EKPSTM_PIXMAP%%@dirrm share/ekpstm
%%PIPS%%share/locale/de/LC_MESSAGES/pips%%PRTYPE%%.mo
%%PIPS_EN%%share/locale/en/LC_MESSAGES/pips%%PRTYPE%%.mo
%%PIPS%%share/locale/es/LC_MESSAGES/pips%%PRTYPE%%.mo
%%PIPS%%share/locale/fr/LC_MESSAGES/pips%%PRTYPE%%.mo
%%PIPS%%share/locale/it/LC_MESSAGES/pips%%PRTYPE%%.mo
%%PIPS%%share/locale/ja/LC_MESSAGES/pips%%PRTYPE%%.mo
%%PIPS%%share/locale/ko/LC_MESSAGES/pips%%PRTYPE%%.mo
%%PIPS%%share/locale/nl/LC_MESSAGES/pips%%PRTYPE%%.mo
%%PIPS%%share/locale/pt/LC_MESSAGES/pips%%PRTYPE%%.mo
%%PIPS%%share/locale/zh/LC_MESSAGES/pips%%PRTYPE%%.mo
%%PIPS%%share/locale/zh_TW/LC_MESSAGES/pips%%PRTYPE%%.mo
%%EKPNAVI_MO%%share/locale/de/LC_MESSAGES/ekpnavi.mo
%%EKPNAVI_MO%%share/locale/es/LC_MESSAGES/ekpnavi.mo
%%EKPNAVI_MO%%share/locale/fr/LC_MESSAGES/ekpnavi.mo
%%EKPNAVI_MO%%share/locale/it/LC_MESSAGES/ekpnavi.mo
%%EKPNAVI_MO%%share/locale/ja/LC_MESSAGES/ekpnavi.mo
%%EKPNAVI_MO%%share/locale/ko/LC_MESSAGES/ekpnavi.mo
%%EKPNAVI_MO%%share/locale/nl/LC_MESSAGES/ekpnavi.mo
%%EKPNAVI_MO%%share/locale/pt/LC_MESSAGES/ekpnavi.mo
%%EKPNAVI_MO%%share/locale/zh/LC_MESSAGES/ekpnavi.mo
%%EKPNAVI_MO%%share/locale/zh_TW/LC_MESSAGES/ekpnavi.mo
%%EKPSTM_MO%%share/locale/de/LC_MESSAGES/ekpstm.mo
%%EKPSTM_MO%%share/locale/es/LC_MESSAGES/ekpstm.mo
%%EKPSTM_MO%%share/locale/fr/LC_MESSAGES/ekpstm.mo
%%EKPSTM_MO%%share/locale/it/LC_MESSAGES/ekpstm.mo
%%EKPSTM_MO%%share/locale/ja/LC_MESSAGES/ekpstm.mo
%%EKPSTM_MO%%share/locale/ko/LC_MESSAGES/ekpstm.mo
%%EKPSTM_MO%%share/locale/nl/LC_MESSAGES/ekpstm.mo
%%EKPSTM_MO%%share/locale/pt/LC_MESSAGES/ekpstm.mo
%%EKPSTM_MO%%share/locale/zh/LC_MESSAGES/ekpstm.mo
%%EKPSTM_MO%%share/locale/zh_TW/LC_MESSAGES/ekpstm.mo
