# security/crowdsec-firewall-bouncer/files/crowdsec_firewall.in
#!/bin/sh
#
# PROVIDE: crowdsec_firewall
# REQUIRE: LOGIN DAEMON NETWORKING crowdsec
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# crowdsec_firewall_enable (bool):	Set it to YES to enable crowdsec firewall.
#					Default is "NO"

. /etc/rc.subr

name=crowdsec_firewall
desc="Crowdsec Firewall"
rcvar=crowdsec_firewall_enable

load_rc_config $name

: "${crowdsec_firewall_enable:=NO}"
: "${crowdsec_firewall_config:=%%ETCDIR%%/crowdsec-firewall-bouncer.yaml}"

pidfile=/var/run/${name}.pid
required_files="$crowdsec_firewall_config"
command="%%PREFIX%%/bin/crowdsec-firewall-bouncer"
start_cmd="${name}_start"
start_precmd="${name}_precmd"

crowdsec_firewall_precmd() {
    CSCLI=%%PREFIX%%/bin/cscli
    orig_line="api_key: \${API_KEY}"
    if grep -q "^${orig_line}" "${crowdsec_firewall_config}"; then
        SUFFIX=$(LC_CTYPE=C tr -dc A-Za-z0-9 </dev/urandom | head -c 8)
        BOUNCER="cs-firewall-bouncer-${SUFFIX}"
        if command -v "$CSCLI" >/dev/null; then
            API_KEY=$($CSCLI bouncers add "${BOUNCER}" -o raw)
            if [ -n "$API_KEY" ]; then
                sed -i "" "s/^${orig_line}/api_key: ${API_KEY}     # ${BOUNCER}/" "${crowdsec_firewall_config}"
                echo "Registered: ${BOUNCER}"
            fi
        fi
    fi

    # needs real tabs
    cat <<-EOT | /sbin/pfctl -f /dev/fd/0
	table <crowdsec-blacklists> persist
	table <crowdsec6-blacklists> persist
	block drop in quick from <crowdsec-blacklists> to any
	block drop in quick from <crowdsec6-blacklists> to any
	EOT

}

crowdsec_firewall_start() {
    /usr/sbin/daemon -f -p ${pidfile} -t "${desc}" -- \
        ${command} -c "${crowdsec_firewall_config}"
}

run_rc_command "$1"
