# security/crowdsec-firewall-bouncer/files/crowdsec_firewall.in
#!/bin/sh
#
# PROVIDE: crowdsec_firewall
# REQUIRE: LOGIN DAEMON NETWORKING
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# crowdsec_firewall_enable (bool):	Set it to YES to enable crowdsec firewall.
#					Default is "NO"

. /etc/rc.subr

name=crowdsec_firewall
desc="Crowdsec Firewall"
rcvar=crowdsec_firewall_enable

load_rc_config $name

: ${crowdsec_firewall_enable:="NO"}
: ${crowdsec_firewall_config:="%%PREFIX%%/etc/crowdsec-firewall-bouncer/crowdsec-firewall-bouncer.yaml"}

pidfile=/var/run/$name.pid
required_files="$crowdsec_firewall_config"
procname="%%PREFIX%%/bin/crowdsec-firewall-bouncer"
command=/usr/sbin/daemon
command_args="-fp $pidfile -t '$desc' -- '$procname' -c '$crowdsec_firewall_config'"

run_rc_command "$1"
