# net/kafka/files/kafka.in
#!/bin/sh
# $FreeBSD$
#
# PROVIDE: kafka
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add kafka_enable="YES" to /etc/rc.conf to enable Kafka:
#
# Additional variables you can define are:
#
# kafka_user:		Username to run Kafka
#			Default: %%USERS%%
# kafka_group:		Groupname to run Kafka
#			Default: %%GROUPS%%
# kafka_config:		Configuration file to run Kafka
#			Default: %%ETCDIR%%/server.properties
# kafka_log4j_profile:	Profile for Kafka logging
#			Default: %%ETCDIR%%/log4j.properties
# kafka_log_dir:	Directory to store Kafka logs
#			Default: %%KAFKA_LOGDIR%%
# kafka_java_opts:	Options passed to JVM to start Kafka
#			Default: None
# kafka_pid_file:	Full path of the Kafka process PID file
#			Default: %%KAFKA_RUNDIR%%/kafka.pid
# kafka_out_file:	Filename to store Kafka's stdout and stderr
#			Default: kafkaServer.out

. /etc/rc.subr

name=kafka
rcvar=kafka_enable

load_rc_config "${name}"

: ${kafka_enable:=NO}
: ${kafka_user:=%%USERS%%}
: ${kafka_group:=%%GROUPS%%}
: ${kafka_config:=%%ETCDIR%%/server.properties}
: ${kafka_log4j_profile:=%%ETCDIR%%/log4j.properties}
: ${kafka_log_dir:=%%KAFKA_LOGDIR%%}
: ${kafka_java_opts:=}
: ${kafka_pid_file:=%%KAFKA_RUNDIR%%/kafka.pid}
: ${kafka_out_file:=kafkaServer.out}

JAVA=%%JAVA%%

CLASSPATH=":%%DATADIR%%/libs/*"

kafka_log_opts="-Dkafka.logs.dir=${kafka_log_dir} -Dlog4j.configuration=file:${kafka_log4j_profile}"
kafka_main="$kafka_java_opts $kafka_log_opts -cp $CLASSPATH"
pidfile=${kafka_pid_file}
procname=${JAVA}
required_files="${kafka_config} ${kafka_log4j_profile}"

daemon_head="-f"
daemon_tail=""

if [ -n "$kafka_out_file" ]; then
	daemon_head=""
	daemon_tail=" > ${kafka_log_dir}/${kafka_out_file} 2>&1 < /dev/null"
	[ -d "$kafka_log_dir" ] || mkdir -p "$kafka_log_dir"
	[ -f "${kafka_log_dir}/${kafka_out_file}" ] || touch "${kafka_log_dir}/${kafka_out_file}"
fi

command="/usr/sbin/daemon"
command_args="${daemon_head} -p ${pidfile} ${JAVA} ${kafka_main} kafka.Kafka ${kafka_config} ${daemon_tail}"

start_precmd="kafka_init"

kafka_init()
{
	chown ${kafka_user}:${kafka_group} ${pidfile%/*} ${kafka_log_dir} ${kafka_log_dir}/${kafka_out_file} ${kafka_config} ${kafka_log4j_profile}
}

run_rc_command "$1"
