# mail/prayer/pkg-plist
@exec if ! pw groupshow prayer; then pw groupadd prayer; fi
@exec if ! pw usershow prayer; then pw useradd prayer -g prayer -h - -d /nonexistent -s /nonexistent -c "Prayer Webmail"; fi
etc/rc.d/prayer.sh.sample
prayer/help/CVS/Root
prayer/help/CVS/Repository
prayer/help/CVS/Entries
prayer/help/CVS/Entries.Log
prayer/help/compose/CVS/Root
prayer/help/compose/CVS/Repository
prayer/help/compose/CVS/Entries
prayer/help/compose/large.html
prayer/help/compose/normal.html
prayer/help/compose/postponed.html
prayer/help/compose/role_list.html
prayer/help/abook_add.html
prayer/help/abook_list.html
prayer/help/abook_search.html
prayer/help/abook_xfer.html
prayer/help/attachments.html
prayer/help/block.html
prayer/help/copy.html
prayer/help/dictionary.html
prayer/help/display.html
prayer/help/favourites.html
prayer/help/filter.html
prayer/help/filter_select.html
prayer/help/folders.html
prayer/help/fullname.html
prayer/help/include.html
prayer/help/list.html
prayer/help/manage.html
prayer/help/passwd.html
prayer/help/quota.html
prayer/help/rename.html
prayer/help/redirect.html
prayer/help/reply.html
prayer/help/roles_entry.html
prayer/help/roles_list.html
prayer/help/search.html
prayer/help/sieve.html
prayer/help/spam.html
prayer/help/toolbar.html
prayer/help/transfer.html
prayer/help/upload_select.html
prayer/help/vacation.html
prayer/help/vaclog.html
prayer/help/welcome.html
prayer/help/prefs/CVS/Root
prayer/help/prefs/CVS/Repository
prayer/help/prefs/CVS/Entries
prayer/help/prefs/compose1.html
prayer/help/prefs/compose2.html
prayer/help/prefs/display1.html
prayer/help/prefs/folder.html
prayer/help/prefs/general.html
prayer/icons/CVS/Root
prayer/icons/CVS/Repository
prayer/icons/CVS/Entries
prayer/icons/addressbook-inactive.gif
prayer/icons/addressbook.gif
prayer/icons/back.gif
prayer/icons/compose-inactive.gif
prayer/icons/compose.gif
prayer/icons/cross.gif
prayer/icons/dir.gif
prayer/icons/display-inactive.gif
prayer/icons/display.gif
prayer/icons/folders-inactive.gif
prayer/icons/folders.gif
prayer/icons/help-inactive.gif
prayer/icons/help.gif
prayer/icons/left-inactive.gif
prayer/icons/left.gif
prayer/icons/leftleft-inactive.gif
prayer/icons/leftleft.gif
prayer/icons/logout.gif
prayer/icons/mailbox-inactive.gif
prayer/icons/mailbox.gif
prayer/icons/manage-inactive.gif
prayer/icons/manage.gif
prayer/icons/msg-answered.gif
prayer/icons/msg-deleted.gif
prayer/icons/msg-seen.gif
prayer/icons/msg-unseen.gif
prayer/icons/off.gif
prayer/icons/on.gif
prayer/icons/paperclip-small.gif
prayer/icons/right-inactive.gif
prayer/icons/right.gif
prayer/icons/rightright-inactive.gif
prayer/icons/rightright.gif
prayer/icons/tick.gif
prayer/icons/univheader.gif
prayer/etc/prayer.cf-dist
prayer/etc/prayer-accountd.cf-dist
prayer/etc/motd.html
prayer/etc/welcome.html
prayer/sbin/prayer
prayer/sbin/prayer-session
prayer/sbin/prayer-ssl-prune
prayer/sbin/prayer-sem-prune
prayer/sbin/prayer-cyclog
@dirrm prayer/help/CVS
@dirrm prayer/help/compose/CVS
@dirrm prayer/help/compose
@dirrm prayer/help/prefs/CVS
@dirrm prayer/help/prefs
@dirrm prayer/help
@dirrm prayer/icons/CVS
@dirrm prayer/icons
@dirrm prayer/etc
@dirrm prayer/certs
@dirrm prayer/sbin
@dirrm prayer
%%PORTDOCS%%%%DOCSDIR%%/ACCOUNTD
%%PORTDOCS%%%%DOCSDIR%%/CMD_LINE
%%PORTDOCS%%%%DOCSDIR%%/DESIGN
%%PORTDOCS%%%%DOCSDIR%%/DONE
%%PORTDOCS%%%%DOCSDIR%%/DONE-OLD
%%PORTDOCS%%%%DOCSDIR%%/FEATURES
%%PORTDOCS%%%%DOCSDIR%%/ICONLIST
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/LOGS
%%PORTDOCS%%%%DOCSDIR%%/NOTICE
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/ROADMAP
%%PORTDOCS%%%%DOCSDIR%%/SECURITY
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%%%DOCSDIR%%/URL_OPTIONS
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@exec mkdir -p /var/spool/prayer
@exec chown -R prayer:prayer /var/spool/prayer
@unexec rmdir /var/spool/prayer
