# net/google-daemon/files/google-accounts-manager.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: google-accounts-manager
# BEFORE: LOGIN
# REQUIRE: syslogd
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable google-accounts-manager:
#
# google_accounts_manager_enable="YES"
#

. /etc/rc.subr

name="google_accounts_manager"
rcvar="google_accounts_manager_enable"

command_interpreter="%%PREFIX%%/bin/python"
command=%%PREFIX%%/share/google/google_daemon/manage_accounts.py
pidfile=/var/run/manage_accounts.pid
command_args=--daemon

google_accounts_manager_enable=${google_accounts_manager_enable:-"NO"}

load_rc_config "${name}"
run_rc_command "$1"
