# mail/drac/pkg-plist
sbin/rpc.dracd
include/drac.h
lib/libdrac.a
etc/rc.d/dracd.sh.sample
