# net/howl/pkg-plist
bin/mDNSBrowse
bin/mDNSPublish
bin/mDNSQuery
bin/mDNSResolve
bin/mDNSResponder
etc/rc.d/mdnsresponder.sh
include/howl/corby/buffer.h
include/howl/corby/channel.h
include/howl/corby/corby.h
include/howl/corby/message.h
include/howl/corby/object.h
include/howl/corby/orb.h
include/howl/discovery/discovery.h
include/howl/discovery/text_record.h
include/howl/howl.h
include/howl/howl_config.h
include/howl/rendezvous/rendezvous.h
include/howl/rendezvous/text_record.h
include/howl/salt/address.h
include/howl/salt/debug.h
include/howl/salt/interface.h
include/howl/salt/platform.h
include/howl/salt/salt.h
include/howl/salt/signal.h
include/howl/salt/socket.h
include/howl/salt/time.h
lib/libhowl.a
lib/libhowl.so
lib/libhowl.so.0
lib/libmDNSResponder.a
lib/libmDNSResponder.so
lib/libmDNSResponder.so.0
libdata/pkgconfig/howl.pc
%%DATADIR%%/help/address-summary.html
%%DATADIR%%/help/address.html
%%DATADIR%%/help/discovery-summary.html
%%DATADIR%%/help/discovery.html
%%DATADIR%%/help/index.html
%%DATADIR%%/help/overall.html
%%DATADIR%%/help/salt-summary.html
%%DATADIR%%/help/salt.html
%%DATADIR%%/help/txt_rec-summary.html
%%DATADIR%%/help/txt_rec.html
%%DATADIR%%/help/txt_rec_iter-summary.html
%%DATADIR%%/help/txt_rec_iter.html
%%DATADIR%%/help/txt_rec_str_iter-summary.html
%%DATADIR%%/help/txt_rec_str_iter.html
@dirrm %%DATADIR%%/help
@dirrm %%DATADIR%%
@dirrm include/howl/salt
@dirrm include/howl/rendezvous
@dirrm include/howl/discovery
@dirrm include/howl/corby
@dirrm include/howl
