# sysutils/ezjail/pkg-plist
@unexec cmp -s %D/etc/ezjail.conf.sample %D/etc/ezjail.conf && rm -f %D/etc/ezjail.conf || true
etc/ezjail.conf.sample
@exec [ -f %B/ezjail.conf ] || cp -p %D/%F %B/ezjail.conf
etc/rc.d/ezjail
bin/ezjail-admin
share/zsh/site-functions/ezjail-admin
%%EXAMPLESDIR%%/example/etc/make.conf
%%EXAMPLESDIR%%/example/etc/periodic.conf
%%EXAMPLESDIR%%/example/etc/rc.conf
%%EXAMPLESDIR%%/example/etc/rc.d/ezjail.flavour.example
%%EXAMPLESDIR%%/example/usr/local/etc/sudoers
%%EXAMPLESDIR%%/nullmailer-example/ezjail.flavour
%%EXAMPLESDIR%%/nullmailer-example/etc/rc.conf
%%EXAMPLESDIR%%/nullmailer-example/etc/mail/mailer.conf
%%EXAMPLESDIR%%/nullmailer-example/usr/local/etc/nullmailer/remotes

@dirrm etc/ezjail/
@dirrm %%EXAMPLESDIR%%/example/usr/local/etc/
@dirrm %%EXAMPLESDIR%%/example/usr/local/
@dirrm %%EXAMPLESDIR%%/example/usr/
@dirrm %%EXAMPLESDIR%%/example/etc/rc.d/
@dirrm %%EXAMPLESDIR%%/example/etc/
@dirrm %%EXAMPLESDIR%%/example/
@dirrm %%EXAMPLESDIR%%/nullmailer-example/usr/local/etc/nullmailer/
@dirrm %%EXAMPLESDIR%%/nullmailer-example/usr/local/etc/
@dirrm %%EXAMPLESDIR%%/nullmailer-example/usr/local/
@dirrm %%EXAMPLESDIR%%/nullmailer-example/usr/
@dirrm %%EXAMPLESDIR%%/nullmailer-example/etc/mail/
@dirrm %%EXAMPLESDIR%%/nullmailer-example/etc/
@dirrm %%EXAMPLESDIR%%/nullmailer-example/
@dirrm %%EXAMPLESDIR%%
@dirrmtry share/zsh/site-functions
@dirrmtry share/zsh
