# sysutils/ezjail/pkg-plist
etc/ezjail.conf.sample
etc/rc.d/ezjail.sh
bin/ezjail-admin
@dirrm etc/ezjail
