# audio/teamspeak3-server/pkg-plist
etc/rc.d/teamspeak_server.sh
lib/teamspeak_server/changelog.txt
lib/teamspeak_server/httpdocs/choice_box.html
lib/teamspeak_server/httpdocs/database_client_manager.html
lib/teamspeak_server/httpdocs/database_client_manager_add.html
lib/teamspeak_server/httpdocs/database_client_manager_clientlist.html
lib/teamspeak_server/httpdocs/database_client_manager_edit.html
lib/teamspeak_server/httpdocs/database_sclient_manager.html
lib/teamspeak_server/httpdocs/database_sclient_manager_add.html
lib/teamspeak_server/httpdocs/database_sclient_manager_clientlist.html
lib/teamspeak_server/httpdocs/database_sclient_manager_edit.html
lib/teamspeak_server/httpdocs/error_box.html
lib/teamspeak_server/httpdocs/footer.html
lib/teamspeak_server/httpdocs/gfx/add_client.png
lib/teamspeak_server/httpdocs/gfx/add_server.png
lib/teamspeak_server/httpdocs/gfx/delete.png
lib/teamspeak_server/httpdocs/gfx/edit.png
lib/teamspeak_server/httpdocs/gfx/headerline.gif
lib/teamspeak_server/httpdocs/gfx/menu_background.gif
lib/teamspeak_server/httpdocs/gfx/pixel.gif
lib/teamspeak_server/httpdocs/gfx/select.png
lib/teamspeak_server/httpdocs/gfx/select_short.png
lib/teamspeak_server/httpdocs/gfx/server-run.png
lib/teamspeak_server/httpdocs/gfx/server-stop.png
lib/teamspeak_server/httpdocs/gfx/server_icon.png
lib/teamspeak_server/httpdocs/gfx/serveradminlogin.gif
lib/teamspeak_server/httpdocs/gfx/uhr_icon.png
lib/teamspeak_server/httpdocs/gfx/under_hl.gif
lib/teamspeak_server/httpdocs/gfx/user_icon.png
lib/teamspeak_server/httpdocs/gfx/vBulletin_logo.gif
lib/teamspeak_server/httpdocs/header.html
lib/teamspeak_server/httpdocs/index.html
lib/teamspeak_server/httpdocs/login.html
lib/teamspeak_server/httpdocs/login_error.html
lib/teamspeak_server/httpdocs/menu_bottom.html
lib/teamspeak_server/httpdocs/menu_top.html
lib/teamspeak_server/httpdocs/ok_box.html
lib/teamspeak_server/httpdocs/server_basic_settings.html
lib/teamspeak_server/httpdocs/server_country_option.html
lib/teamspeak_server/httpdocs/server_manager.html
lib/teamspeak_server/httpdocs/server_manager_add.html
lib/teamspeak_server/httpdocs/server_manager_grouppermission.html
lib/teamspeak_server/httpdocs/server_manager_overview.html
lib/teamspeak_server/httpdocs/server_manager_permission.html
lib/teamspeak_server/httpdocs/server_manager_permission_ca.html
lib/teamspeak_server/httpdocs/server_manager_permission_op.html
lib/teamspeak_server/httpdocs/server_manager_permission_r.html
lib/teamspeak_server/httpdocs/server_manager_permission_sa.html
lib/teamspeak_server/httpdocs/server_manager_permission_u.html
lib/teamspeak_server/httpdocs/server_manager_permission_v.html
lib/teamspeak_server/httpdocs/server_manager_permissionfield.html
lib/teamspeak_server/httpdocs/server_manager_serverlist.html
lib/teamspeak_server/httpdocs/server_manager_servers.html
lib/teamspeak_server/httpdocs/server_manager_servers_list.html
lib/teamspeak_server/httpdocs/server_manager_settings.html
lib/teamspeak_server/httpdocs/slogin.html
lib/teamspeak_server/readme.txt
lib/teamspeak_server/server_linux
lib/teamspeak_server/slicense.txt
lib/teamspeak_server/sqlite.so
lib/teamspeak_server/sqlite_sql/create_bans.sql
lib/teamspeak_server/sqlite_sql/create_channel_privileges.sql
lib/teamspeak_server/sqlite_sql/create_channels.sql
lib/teamspeak_server/sqlite_sql/create_clients.sql
lib/teamspeak_server/sqlite_sql/create_server_privileges.sql
lib/teamspeak_server/sqlite_sql/create_servers.sql
lib/teamspeak_server/sqlite_sql/delete_ban.sql
lib/teamspeak_server/sqlite_sql/delete_ban_ip.sql
lib/teamspeak_server/sqlite_sql/delete_channel.sql
lib/teamspeak_server/sqlite_sql/delete_client.sql
lib/teamspeak_server/sqlite_sql/delete_server.sql
lib/teamspeak_server/sqlite_sql/edit_channel.sql
lib/teamspeak_server/sqlite_sql/edit_channel_privilege.sql
lib/teamspeak_server/sqlite_sql/edit_client.sql
lib/teamspeak_server/sqlite_sql/edit_server.sql
lib/teamspeak_server/sqlite_sql/edit_server_privilege.sql
lib/teamspeak_server/sqlite_sql/new_ban.sql
lib/teamspeak_server/sqlite_sql/new_channel.sql
lib/teamspeak_server/sqlite_sql/new_channel_privilege.sql
lib/teamspeak_server/sqlite_sql/new_client.sql
lib/teamspeak_server/sqlite_sql/new_server.sql
lib/teamspeak_server/sqlite_sql/new_server_privilege.sql
lib/teamspeak_server/sqlite_sql/read_active_servers.sql
lib/teamspeak_server/sqlite_sql/read_bans.sql
lib/teamspeak_server/sqlite_sql/read_channel_privilege.sql
lib/teamspeak_server/sqlite_sql/read_channels.sql
lib/teamspeak_server/sqlite_sql/read_client_login.sql
lib/teamspeak_server/sqlite_sql/read_clients.sql
lib/teamspeak_server/sqlite_sql/read_clients_id.sql
lib/teamspeak_server/sqlite_sql/read_clients_name.sql
lib/teamspeak_server/sqlite_sql/read_server.sql
lib/teamspeak_server/sqlite_sql/read_server_port.sql
lib/teamspeak_server/sqlite_sql/read_server_privilege.sql
lib/teamspeak_server/sqlite_sql/read_servers.sql
lib/teamspeak_server/sqlite_sql/remove_moderatedattribs.sql
lib/teamspeak_server/sqlite_sql/upgrade_1.sql
lib/teamspeak_server/sqlite_sql/upgrade_2.sql
lib/teamspeak_server/tss.SuSE_startscript
lib/teamspeak_server/tss.minimal_startscript
@dirrm lib/teamspeak_server/sqlite_sql
@dirrm lib/teamspeak_server/httpdocs/gfx
@dirrm lib/teamspeak_server/httpdocs
@dirrm lib/teamspeak_server
