# audio/teamspeak3-server/pkg-plist
etc/rc.d/teamspeak_server.sh.sample
lib/teamspeak_server/INSTALL
lib/teamspeak_server/INSTALL.mysql
lib/teamspeak_server/Manual/C1H_HTML.css
lib/teamspeak_server/Manual/DocToHelp.js
lib/teamspeak_server/Manual/RC2servermanual-blank.htm
lib/teamspeak_server/Manual/RC2servermanual-css.rtf
lib/teamspeak_server/Manual/RC2servermanual-index.htm
lib/teamspeak_server/Manual/RC2servermanual-toc.htm
lib/teamspeak_server/Manual/RC2servermanuallinux
lib/teamspeak_server/Manual/_TeamSpeak RC2 Server manual-1.png
lib/teamspeak_server/Manual/administration1.htm
lib/teamspeak_server/Manual/back.gif
lib/teamspeak_server/Manual/bandwidthandcodecs1.htm
lib/teamspeak_server/Manual/button.gif
lib/teamspeak_server/Manual/closed.gif
lib/teamspeak_server/Manual/closingdowntheserver4.htm
lib/teamspeak_server/Manual/closingdowntheserver5.htm
lib/teamspeak_server/Manual/default.htm
lib/teamspeak_server/Manual/frequentlyaskedquestions1.htm
lib/teamspeak_server/Manual/howtoruntheteamspeakserverasaserviceunderwinxpand20001.htm
lib/teamspeak_server/Manual/installingtheserver4.htm
lib/teamspeak_server/Manual/installingtheserver5.htm
lib/teamspeak_server/Manual/lastnotes1.htm
lib/teamspeak_server/Manual/licensing1.htm
lib/teamspeak_server/Manual/linuxinstalling.htm
lib/teamspeak_server/Manual/loggingin1.htm
lib/teamspeak_server/Manual/login.png
lib/teamspeak_server/Manual/logout1.htm
lib/teamspeak_server/Manual/mainconfig1.htm
lib/teamspeak_server/Manual/miscellaneous1.htm
lib/teamspeak_server/Manual/open.gif
lib/teamspeak_server/Manual/ports1.htm
lib/teamspeak_server/Manual/qilostmyserveradminpassword.htm
lib/teamspeak_server/Manual/qwhatportsdoestheserveruse.htm
lib/teamspeak_server/Manual/qwhatsupportdoesteamspeakoffer.htm
lib/teamspeak_server/Manual/requirementstoruntheteamspeak2server1.htm
lib/teamspeak_server/Manual/serverpermissions1.htm
lib/teamspeak_server/Manual/serversettings.png
lib/teamspeak_server/Manual/serversettings1.htm
lib/teamspeak_server/Manual/startingtheserver1.htm
lib/teamspeak_server/Manual/startinguptheserver1.htm
lib/teamspeak_server/Manual/superadmin1.htm
lib/teamspeak_server/Manual/superadminservers.png
lib/teamspeak_server/Manual/superadminstart.png
lib/teamspeak_server/Manual/tcpqueryfunction1.htm
lib/teamspeak_server/Manual/teamspeak.gif
lib/teamspeak_server/Manual/teamspeakrc2servermanual1.htm
lib/teamspeak_server/Manual/toc.css
lib/teamspeak_server/Manual/topic.gif
lib/teamspeak_server/Manual/usermanager.png
lib/teamspeak_server/Manual/usermanager1.htm
lib/teamspeak_server/Manual/webadminningtheserver1.htm
lib/teamspeak_server/Manual/webpost2.htm
lib/teamspeak_server/Manual/webpost3.htm
lib/teamspeak_server/Manual/whatistheteamspeak2server1.htm
lib/teamspeak_server/Manual/whatsnewinteamspeakreleasecandidate21.htm
lib/teamspeak_server/Manual/windowsinstalling.htm
lib/teamspeak_server/Manual/yourips1.htm
lib/teamspeak_server/README
lib/teamspeak_server/changelog.txt
lib/teamspeak_server/httpdocs/choice_box.html
lib/teamspeak_server/httpdocs/database_client_manager.html
lib/teamspeak_server/httpdocs/database_client_manager_add.html
lib/teamspeak_server/httpdocs/database_client_manager_clientlist.html
lib/teamspeak_server/httpdocs/database_client_manager_edit.html
lib/teamspeak_server/httpdocs/database_sclient_manager.html
lib/teamspeak_server/httpdocs/database_sclient_manager_add.html
lib/teamspeak_server/httpdocs/database_sclient_manager_clientlist.html
lib/teamspeak_server/httpdocs/database_sclient_manager_edit.html
lib/teamspeak_server/httpdocs/error_box.html
lib/teamspeak_server/httpdocs/footer.html
lib/teamspeak_server/httpdocs/gfx/add_client.png
lib/teamspeak_server/httpdocs/gfx/add_server.png
lib/teamspeak_server/httpdocs/gfx/delete.png
lib/teamspeak_server/httpdocs/gfx/edit.png
lib/teamspeak_server/httpdocs/gfx/global_settings.jpg
lib/teamspeak_server/httpdocs/gfx/headerline.gif
lib/teamspeak_server/httpdocs/gfx/help.gif
lib/teamspeak_server/httpdocs/gfx/helpnav.jpg
lib/teamspeak_server/httpdocs/gfx/helppages.jpg
lib/teamspeak_server/httpdocs/gfx/helppages.png
lib/teamspeak_server/httpdocs/gfx/menu_background.gif
lib/teamspeak_server/httpdocs/gfx/permissions.jpg
lib/teamspeak_server/httpdocs/gfx/pixel.gif
lib/teamspeak_server/httpdocs/gfx/select.png
lib/teamspeak_server/httpdocs/gfx/select_short.png
lib/teamspeak_server/httpdocs/gfx/server-run.png
lib/teamspeak_server/httpdocs/gfx/server-stop.png
lib/teamspeak_server/httpdocs/gfx/server_icon.png
lib/teamspeak_server/httpdocs/gfx/server_overview.jpg
lib/teamspeak_server/httpdocs/gfx/server_settings.jpg
lib/teamspeak_server/httpdocs/gfx/serveradminlogin.gif
lib/teamspeak_server/httpdocs/gfx/servers.jpg
lib/teamspeak_server/httpdocs/gfx/superadmin_manager.jpg
lib/teamspeak_server/httpdocs/gfx/uhr_icon.png
lib/teamspeak_server/httpdocs/gfx/under_hl.gif
lib/teamspeak_server/httpdocs/gfx/user_icon.png
lib/teamspeak_server/httpdocs/gfx/user_manager.jpg
lib/teamspeak_server/httpdocs/gfx/vBulletin_logo.gif
lib/teamspeak_server/httpdocs/header.html
lib/teamspeak_server/httpdocs/help/global_settings.html
lib/teamspeak_server/httpdocs/help/permissions_general.html
lib/teamspeak_server/httpdocs/help/server_overview.html
lib/teamspeak_server/httpdocs/help/server_settings.html
lib/teamspeak_server/httpdocs/help/servers.html
lib/teamspeak_server/httpdocs/help/superadmin_manager.html
lib/teamspeak_server/httpdocs/help/user_manager.html
lib/teamspeak_server/httpdocs/index.html
lib/teamspeak_server/httpdocs/login.html
lib/teamspeak_server/httpdocs/login_error.html
lib/teamspeak_server/httpdocs/menu_bottom.html
lib/teamspeak_server/httpdocs/menu_top.html
lib/teamspeak_server/httpdocs/ok_box.html
lib/teamspeak_server/httpdocs/server_basic_settings.html
lib/teamspeak_server/httpdocs/server_country_option.html
lib/teamspeak_server/httpdocs/server_manager.html
lib/teamspeak_server/httpdocs/server_manager_add.html
lib/teamspeak_server/httpdocs/server_manager_grouppermission.html
lib/teamspeak_server/httpdocs/server_manager_overview.html
lib/teamspeak_server/httpdocs/server_manager_permission.html
lib/teamspeak_server/httpdocs/server_manager_permission_ca.html
lib/teamspeak_server/httpdocs/server_manager_permission_op.html
lib/teamspeak_server/httpdocs/server_manager_permission_r.html
lib/teamspeak_server/httpdocs/server_manager_permission_sa.html
lib/teamspeak_server/httpdocs/server_manager_permission_u.html
lib/teamspeak_server/httpdocs/server_manager_permission_v.html
lib/teamspeak_server/httpdocs/server_manager_permissionfield.html
lib/teamspeak_server/httpdocs/server_manager_serverlist.html
lib/teamspeak_server/httpdocs/server_manager_servers.html
lib/teamspeak_server/httpdocs/server_manager_servers_list.html
lib/teamspeak_server/httpdocs/server_manager_settings.html
lib/teamspeak_server/httpdocs/slogin.html
lib/teamspeak_server/libsqlmy.so
lib/teamspeak_server/manual.html
lib/teamspeak_server/mysql_sql/create_bans.sql
lib/teamspeak_server/mysql_sql/create_channel_privileges.sql
lib/teamspeak_server/mysql_sql/create_channels.sql
lib/teamspeak_server/mysql_sql/create_clients.sql
lib/teamspeak_server/mysql_sql/create_nicks.sql
lib/teamspeak_server/mysql_sql/create_server_privileges.sql
lib/teamspeak_server/mysql_sql/create_servers.sql
lib/teamspeak_server/mysql_sql/create_settings.sql
lib/teamspeak_server/mysql_sql/delete_ban.sql
lib/teamspeak_server/mysql_sql/delete_ban_ip.sql
lib/teamspeak_server/mysql_sql/delete_channel.sql
lib/teamspeak_server/mysql_sql/delete_channel_privilege.sql
lib/teamspeak_server/mysql_sql/delete_client.sql
lib/teamspeak_server/mysql_sql/delete_server.sql
lib/teamspeak_server/mysql_sql/edit_channel.sql
lib/teamspeak_server/mysql_sql/edit_channel_privilege.sql
lib/teamspeak_server/mysql_sql/edit_client.sql
lib/teamspeak_server/mysql_sql/edit_client_lastonline.sql
lib/teamspeak_server/mysql_sql/edit_server.sql
lib/teamspeak_server/mysql_sql/edit_server_privilege.sql
lib/teamspeak_server/mysql_sql/new_ban.sql
lib/teamspeak_server/mysql_sql/new_channel.sql
lib/teamspeak_server/mysql_sql/new_channel_privilege.sql
lib/teamspeak_server/mysql_sql/new_client.sql
lib/teamspeak_server/mysql_sql/new_server.sql
lib/teamspeak_server/mysql_sql/new_server_privilege.sql
lib/teamspeak_server/mysql_sql/read_active_servers.sql
lib/teamspeak_server/mysql_sql/read_bans.sql
lib/teamspeak_server/mysql_sql/read_channel.sql
lib/teamspeak_server/mysql_sql/read_channel_privilege.sql
lib/teamspeak_server/mysql_sql/read_channel_privileges.sql
lib/teamspeak_server/mysql_sql/read_channels.sql
lib/teamspeak_server/mysql_sql/read_client_channels_privileges.sql
lib/teamspeak_server/mysql_sql/read_client_login.sql
lib/teamspeak_server/mysql_sql/read_clients.sql
lib/teamspeak_server/mysql_sql/read_clients_id.sql
lib/teamspeak_server/mysql_sql/read_clients_name.sql
lib/teamspeak_server/mysql_sql/read_server.sql
lib/teamspeak_server/mysql_sql/read_server_port.sql
lib/teamspeak_server/mysql_sql/read_server_privilege.sql
lib/teamspeak_server/mysql_sql/read_servers.sql
lib/teamspeak_server/mysql_sql/remove_bans.sql
lib/teamspeak_server/mysql_sql/remove_moderatedattribs.sql
lib/teamspeak_server/mysql_sql/search_client.sql
lib/teamspeak_server/mysql_sql/upgrade_1.sql
lib/teamspeak_server/mysql_sql/upgrade_2.sql
lib/teamspeak_server/mysql_sql/upgrade_3.sql
lib/teamspeak_server/server_linux
lib/teamspeak_server/sqlite.so
lib/teamspeak_server/sqlite_sql/create_bans.sql
lib/teamspeak_server/sqlite_sql/create_channel_privileges.sql
lib/teamspeak_server/sqlite_sql/create_channels.sql
lib/teamspeak_server/sqlite_sql/create_clients.sql
lib/teamspeak_server/sqlite_sql/create_nicks.sql
lib/teamspeak_server/sqlite_sql/create_server_privileges.sql
lib/teamspeak_server/sqlite_sql/create_servers.sql
lib/teamspeak_server/sqlite_sql/create_settings.sql
lib/teamspeak_server/sqlite_sql/delete_ban.sql
lib/teamspeak_server/sqlite_sql/delete_ban_ip.sql
lib/teamspeak_server/sqlite_sql/delete_channel.sql
lib/teamspeak_server/sqlite_sql/delete_channel_privilege.sql
lib/teamspeak_server/sqlite_sql/delete_client.sql
lib/teamspeak_server/sqlite_sql/delete_server.sql
lib/teamspeak_server/sqlite_sql/edit_channel.sql
lib/teamspeak_server/sqlite_sql/edit_channel_privilege.sql
lib/teamspeak_server/sqlite_sql/edit_client.sql
lib/teamspeak_server/sqlite_sql/edit_client_lastonline.sql
lib/teamspeak_server/sqlite_sql/edit_server.sql
lib/teamspeak_server/sqlite_sql/edit_server_privilege.sql
lib/teamspeak_server/sqlite_sql/new_ban.sql
lib/teamspeak_server/sqlite_sql/new_channel.sql
lib/teamspeak_server/sqlite_sql/new_channel_privilege.sql
lib/teamspeak_server/sqlite_sql/new_client.sql
lib/teamspeak_server/sqlite_sql/new_server.sql
lib/teamspeak_server/sqlite_sql/new_server_privilege.sql
lib/teamspeak_server/sqlite_sql/read_active_servers.sql
lib/teamspeak_server/sqlite_sql/read_bans.sql
lib/teamspeak_server/sqlite_sql/read_channel.sql
lib/teamspeak_server/sqlite_sql/read_channel_privilege.sql
lib/teamspeak_server/sqlite_sql/read_channel_privileges.sql
lib/teamspeak_server/sqlite_sql/read_channels.sql
lib/teamspeak_server/sqlite_sql/read_client_channels_privileges.sql
lib/teamspeak_server/sqlite_sql/read_client_login.sql
lib/teamspeak_server/sqlite_sql/read_clients.sql
lib/teamspeak_server/sqlite_sql/read_clients_id.sql
lib/teamspeak_server/sqlite_sql/read_clients_name.sql
lib/teamspeak_server/sqlite_sql/read_server.sql
lib/teamspeak_server/sqlite_sql/read_server_port.sql
lib/teamspeak_server/sqlite_sql/read_server_privilege.sql
lib/teamspeak_server/sqlite_sql/read_servers.sql
lib/teamspeak_server/sqlite_sql/remove_bans.sql
lib/teamspeak_server/sqlite_sql/remove_moderatedattribs.sql
lib/teamspeak_server/sqlite_sql/search_client.sql
lib/teamspeak_server/sqlite_sql/upgrade_1.sql
lib/teamspeak_server/sqlite_sql/upgrade_2.sql
lib/teamspeak_server/sqlite_sql/upgrade_3.sql
lib/teamspeak_server/tcpquerydocs/banadd.txt
lib/teamspeak_server/tcpquerydocs/banclear.txt
lib/teamspeak_server/tcpquerydocs/bandel.txt
lib/teamspeak_server/tcpquerydocs/banlist.txt
lib/teamspeak_server/tcpquerydocs/banplayer.txt
lib/teamspeak_server/tcpquerydocs/checkserverok.txt
lib/teamspeak_server/tcpquerydocs/ci.txt
lib/teamspeak_server/tcpquerydocs/cl.txt
lib/teamspeak_server/tcpquerydocs/dbci.txt
lib/teamspeak_server/tcpquerydocs/dbfp.txt
lib/teamspeak_server/tcpquerydocs/dbpi.txt
lib/teamspeak_server/tcpquerydocs/dbserverlist.txt
lib/teamspeak_server/tcpquerydocs/dbsuseradd.txt
lib/teamspeak_server/tcpquerydocs/dbsuserchangepw.txt
lib/teamspeak_server/tcpquerydocs/dbsuserdel.txt
lib/teamspeak_server/tcpquerydocs/dbsuserlist.txt
lib/teamspeak_server/tcpquerydocs/dbuseradd.txt
lib/teamspeak_server/tcpquerydocs/dbuserchangeattribs.txt
lib/teamspeak_server/tcpquerydocs/dbuserchangepw.txt
lib/teamspeak_server/tcpquerydocs/dbuserdel.txt
lib/teamspeak_server/tcpquerydocs/dbuserid.txt
lib/teamspeak_server/tcpquerydocs/dbuserlist.txt
lib/teamspeak_server/tcpquerydocs/dccl.txt
lib/teamspeak_server/tcpquerydocs/fc.txt
lib/teamspeak_server/tcpquerydocs/fp.txt
lib/teamspeak_server/tcpquerydocs/gapl.txt
lib/teamspeak_server/tcpquerydocs/gi.txt
lib/teamspeak_server/tcpquerydocs/globalset.txt
lib/teamspeak_server/tcpquerydocs/globalset.txt.bak
lib/teamspeak_server/tcpquerydocs/help.txt
lib/teamspeak_server/tcpquerydocs/ki.txt
lib/teamspeak_server/tcpquerydocs/kick.txt
lib/teamspeak_server/tcpquerydocs/lc.txt
lib/teamspeak_server/tcpquerydocs/log.txt
lib/teamspeak_server/tcpquerydocs/logfind.txt
lib/teamspeak_server/tcpquerydocs/login.txt
lib/teamspeak_server/tcpquerydocs/logmark.txt
lib/teamspeak_server/tcpquerydocs/mptc.txt
lib/teamspeak_server/tcpquerydocs/msg.txt
lib/teamspeak_server/tcpquerydocs/msgall.txt
lib/teamspeak_server/tcpquerydocs/msgu.txt
lib/teamspeak_server/tcpquerydocs/pi.txt
lib/teamspeak_server/tcpquerydocs/pl.txt
lib/teamspeak_server/tcpquerydocs/rehash.txt
lib/teamspeak_server/tcpquerydocs/removeclient.txt
lib/teamspeak_server/tcpquerydocs/sel.txt
lib/teamspeak_server/tcpquerydocs/serveradd.txt
lib/teamspeak_server/tcpquerydocs/serverdel.txt
lib/teamspeak_server/tcpquerydocs/serverset.txt
lib/teamspeak_server/tcpquerydocs/serverstart.txt
lib/teamspeak_server/tcpquerydocs/serverstop.txt
lib/teamspeak_server/tcpquerydocs/si.txt
lib/teamspeak_server/tcpquerydocs/sl.txt
lib/teamspeak_server/tcpquerydocs/slogin.txt
lib/teamspeak_server/tcpquerydocs/sppriv.txt
lib/teamspeak_server/tcpquerydocs/ver.txt
lib/teamspeak_server/teamspeak2-server_startscript
@dirrm lib/teamspeak_server/Manual
@dirrm lib/teamspeak_server/httpdocs/gfx
@dirrm lib/teamspeak_server/httpdocs/help
@dirrm lib/teamspeak_server/httpdocs
@dirrm lib/teamspeak_server/mysql_sql
@dirrm lib/teamspeak_server/sqlite_sql
@dirrm lib/teamspeak_server/tcpquerydocs
@dirrm lib/teamspeak_server
