# sysutils/nut/pkg-plist
sbin/apcsmart
sbin/bestfort
sbin/bestups
sbin/fentonups
sbin/genericups
sbin/ipt-anzen
sbin/mustekups                       
sbin/optiups
sbin/ups-trust425+625
sbin/victronups
sbin/upsc
sbin/upscmd
sbin/upslog
sbin/upsmon
sbin/upsct
sbin/upsct2
sbin/upsd
etc/nut/hosts.conf.sample
etc/nut/multimon.conf.sample
etc/nut/upsd.conf.sample
etc/nut/upsmon.conf.sample
etc/nut/upsset.passwd.sample
etc/rc.d/upsd.sh.sample
@dirrm etc/nut
