# sysutils/nut/pkg-plist
bin/upsc
bin/upscmd
bin/upslog
bin/upsrw
etc/nut/ups.conf.sample
etc/nut/upsd.conf.sample
etc/nut/upsd.users.sample
etc/nut/upsmon.conf.sample
etc/nut/upssched.conf.sample
etc/rc.d/nut.sh.sample
libexec/nut/apcsmart
libexec/nut/bcmxcp
libexec/nut/belkin
libexec/nut/bestferrups801-807
libexec/nut/bestfortress
libexec/nut/bestuferrups
libexec/nut/bestups
libexec/nut/cyberpower
libexec/nut/etapro
libexec/nut/everups
libexec/nut/fentonups
libexec/nut/genericups
libexec/nut/hp
libexec/nut/isbmex
libexec/nut/liebert
libexec/nut/masterguard
libexec/nut/mge-ellipse
libexec/nut/mge-utalk
libexec/nut/microdowell
libexec/nut/newapc
libexec/nut/oneac
libexec/nut/powercom
libexec/nut/sec
libexec/nut/sms
libexec/nut/tripplite
libexec/nut/upsdrvctl
libexec/nut/victronups
sbin/upsd
sbin/upsmon
sbin/upssched
%%PORTDOCS%%share/doc/nut/CHANGES
%%PORTDOCS%%share/doc/nut/CREDITS
%%PORTDOCS%%share/doc/nut/FAQ
%%PORTDOCS%%share/doc/nut/INSTALL
%%PORTDOCS%%share/doc/nut/NEWS
%%PORTDOCS%%share/doc/nut/README
%%PORTDOCS%%share/doc/nut/UPGRADING
%%PORTDOCS%%share/doc/nut/big-servers.txt
%%PORTDOCS%%share/doc/nut/cables/apc.txt
%%PORTDOCS%%share/doc/nut/cables/ge-imv-victron.txt
%%PORTDOCS%%share/doc/nut/cables/imv.txt
%%PORTDOCS%%share/doc/nut/cables/mgeups.txt
%%PORTDOCS%%share/doc/nut/cables/powerware.txt
%%PORTDOCS%%share/doc/nut/cables/repotec.txt
%%PORTDOCS%%share/doc/nut/cables/sms.txt
%%PORTDOCS%%share/doc/nut/chroot.txt
%%PORTDOCS%%share/doc/nut/commands.txt
%%PORTDOCS%%share/doc/nut/config-files.txt
%%PORTDOCS%%share/doc/nut/configure.txt
%%PORTDOCS%%share/doc/nut/data-room.txt
%%PORTDOCS%%share/doc/nut/design.txt
%%PORTDOCS%%share/doc/nut/developers.txt
%%PORTDOCS%%share/doc/nut/driver.list
%%PORTDOCS%%share/doc/nut/drivers/README
%%PORTDOCS%%share/doc/nut/drivers/generic-ups.txt
%%PORTDOCS%%share/doc/nut/drivers/mge-ellipse.txt
%%PORTDOCS%%share/doc/nut/drivers/mge-utalk.txt
%%PORTDOCS%%share/doc/nut/drivers/snmp-ups.txt
%%PORTDOCS%%share/doc/nut/new-drivers.txt
%%PORTDOCS%%share/doc/nut/pager.txt
%%PORTDOCS%%share/doc/nut/protocol.txt
%%PORTDOCS%%share/doc/nut/shutdown.txt
%%PORTDOCS%%share/doc/nut/ssl.txt
%%PORTDOCS%%share/doc/nut/todo.txt
%%PORTDOCS%%share/doc/nut/upssched.txt
%%PORTDOCS%%@dirrm share/doc/nut/cables
%%PORTDOCS%%@dirrm share/doc/nut/drivers
%%PORTDOCS%%@dirrm share/doc/nut
@dirrm etc/nut
@dirrm libexec/nut
