# sysutils/nut/pkg-plist
bin/upsc
bin/upscmd
bin/upsct
bin/upsct2
bin/upslog
etc/nut/hosts.conf.sample
etc/nut/multimon.conf.sample
etc/nut/upsd.conf.sample
etc/nut/upsmon.conf.sample
etc/nut/upsset.passwd.sample
etc/rc.d/nut.sh.sample
libexec/nut/apcsmart
libexec/nut/belkin
libexec/nut/bestfort
libexec/nut/bestuferrups
libexec/nut/bestups
libexec/nut/engetron
libexec/nut/fentonups
libexec/nut/genericups
libexec/nut/ipt-anzen
libexec/nut/mustekups
libexec/nut/newapc
libexec/nut/optiups
libexec/nut/powercom
libexec/nut/toshiba1500
libexec/nut/ups-trust425+625
libexec/nut/victronups
sbin/upsd
sbin/upsmon
sbin/upssched
@dirrm etc/nut
@dirrm libexec/nut
