# sysutils/nut/pkg-plist
bin/upsc
bin/upscmd
bin/upsct
bin/upsct2
bin/upslog
etc/nut/hosts.conf.sample
etc/nut/multimon.conf.sample
etc/nut/upsd.conf.sample
etc/nut/upsmon.conf.sample
etc/nut/upsset.passwd.sample
etc/rc.d/nut.sh.sample
libexec/nut/apcsmart
libexec/nut/belkin
libexec/nut/bestfort
libexec/nut/bestuferrups
libexec/nut/bestups
libexec/nut/engetron
libexec/nut/fentonups
libexec/nut/genericups
libexec/nut/ipt-anzen
libexec/nut/mge-ellipse
libexec/nut/mgeups
libexec/nut/mustekups
libexec/nut/newapc
libexec/nut/optiups
libexec/nut/powercom
libexec/nut/toshiba1500
libexec/nut/ups-trust425+625
libexec/nut/upseyeux
libexec/nut/victronups
sbin/upsd
sbin/upsmon
sbin/upssched
%%PORTDOCS%%share/doc/nut/Changes.trust
%%PORTDOCS%%share/doc/nut/FAQ
%%PORTDOCS%%share/doc/nut/QUICKSTART
%%PORTDOCS%%share/doc/nut/README
%%PORTDOCS%%share/doc/nut/access.txt
%%PORTDOCS%%share/doc/nut/belkin.txt
%%PORTDOCS%%share/doc/nut/big-servers.txt
%%PORTDOCS%%share/doc/nut/commands.txt
%%PORTDOCS%%share/doc/nut/configure.txt
%%PORTDOCS%%share/doc/nut/data-room.txt
%%PORTDOCS%%share/doc/nut/design.txt
%%PORTDOCS%%share/doc/nut/dummyups.txt
%%PORTDOCS%%share/doc/nut/generic-ups.txt
%%PORTDOCS%%share/doc/nut/mge-ellipse.txt
%%PORTDOCS%%share/doc/nut/model-arguments.txt
%%PORTDOCS%%share/doc/nut/new-modules.txt
%%PORTDOCS%%share/doc/nut/pager.txt
%%PORTDOCS%%share/doc/nut/powercom.txt
%%PORTDOCS%%share/doc/nut/protocol.txt
%%PORTDOCS%%share/doc/nut/shutdown.txt
%%PORTDOCS%%share/doc/nut/tips.txt
%%PORTDOCS%%share/doc/nut/todo.txt
%%PORTDOCS%%share/doc/nut/ups-trust425+625.txt
%%PORTDOCS%%share/doc/nut/upssched.txt
%%PORTDOCS%%share/doc/nut/cables/powerware.txt
%%PORTDOCS%%share/doc/nut/cables/victron.txt
%%CGIFILES%%share/nut/cgi/multimon.cgi
%%CGIFILES%%share/nut/cgi/upsimage.cgi
%%CGIFILES%%share/nut/cgi/upsset.cgi
%%CGIFILES%%share/nut/cgi/upsstats.cgi
%%PORTDOCS%%@dirrm share/doc/nut/cables
%%PORTDOCS%%@dirrm share/doc/nut
@dirrm share/nut/cgi
@dirrm share/nut
@dirrm etc/nut
@dirrm libexec/nut
