# sysutils/nut/pkg-plist
@unexec %D/etc/rc.d/nut_upslog.sh forcestop 2>/dev/null || true
@unexec %D/etc/rc.d/nut_upsmon.sh forcestop 2>/dev/null || true
@unexec %D/etc/rc.d/nut.sh forcestop 2>/dev/null || true
@unexec %D/etc/rc.d/nut_upslog forcestop 2>/dev/null || true
@unexec %D/etc/rc.d/nut_upsmon forcestop 2>/dev/null || true
@unexec %D/etc/rc.d/nut forcestop 2>/dev/null || true
%%NUT_LIB%%bin/libupsclient-config
bin/upsc
bin/upscmd
bin/upslog
bin/upsrw
bin/upssched-cmd
etc/nut/cmdvartab
etc/nut/driver.list
etc/nut/ups.conf.sample
etc/nut/upsd.conf.sample
etc/nut/upsd.users.sample
etc/nut/upsmon.conf.sample
etc/nut/upssched.conf.sample
%%NUT_LIB%%include/parseconf.h
%%NUT_LIB%%include/upsclient.h
%%NUT_LIB%%lib/libupsclient.a
libexec/nut/al175
libexec/nut/apcsmart
libexec/nut/bcmxcp
%%NUT_USB%%libexec/nut/bcmxcp_usb
libexec/nut/belkin
libexec/nut/belkinunv
libexec/nut/bestfcom
libexec/nut/bestuferrups
libexec/nut/bestups
libexec/nut/cpsups
libexec/nut/cyberpower
libexec/nut/dummy-ups
libexec/nut/etapro
libexec/nut/everups
libexec/nut/gamatronic
libexec/nut/genericups
libexec/nut/isbmex
libexec/nut/liebert
libexec/nut/masterguard
libexec/nut/megatec
%%NUT_USB%%libexec/nut/megatec_usb
libexec/nut/metasys
libexec/nut/mge-shut
libexec/nut/mge-utalk
libexec/nut/newmge-shut
libexec/nut/nitram
libexec/nut/oneac
libexec/nut/optiups
libexec/nut/powercom
libexec/nut/powerpanel
libexec/nut/rhino
libexec/nut/safenet
%%NUT_SNMP%%libexec/nut/snmp-ups
libexec/nut/skel
libexec/nut/solis
libexec/nut/tripplite
%%NUT_USB%%libexec/nut/tripplite_usb
libexec/nut/tripplitesu
libexec/nut/upscode2
libexec/nut/upsdrvctl
%%NUT_USB%%libexec/nut/usbhid-ups
libexec/nut/victronups
%%NUT_LIB%%libdata/pkgconfig/libupsclient.pc
sbin/upsd
sbin/upsmon
sbin/upssched
%%NUT_CGI%%%%CGIDIR%%/nut/upsimage.cgi
%%NUT_CGI%%%%CGIDIR%%/nut/upsset.cgi
%%NUT_CGI%%%%CGIDIR%%/nut/upsstats.cgi
%%NUT_CGI%%%%CGIETCDIR%%/hosts.conf.sample
%%NUT_CGI%%%%CGIETCDIR%%/upsset.conf.sample
%%NUT_CGI%%%%CGIETCDIR%%/upsstats.html.sample
%%NUT_CGI%%%%CGIETCDIR%%/upsstats-single.html.sample
%%NUT_CGI%%%%WWWDIR%%/nut/bottom.html
%%NUT_CGI%%%%WWWDIR%%/nut/header.html
%%NUT_CGI%%%%WWWDIR%%/nut/index.html
%%NUT_CGI%%%%WWWDIR%%/nut/nut-banner.png
@dirrmtry etc/nut
@dirrmtry lib/pkgconfig
@dirrmtry libexec/nut
@dirrmtry %%STATEDIR%%
%%NUT_CGI%%@dirrm %%CGIDIR%%/nut
%%NUT_CGI%%@dirrmtry %%CGIDIR%%
%%NUT_CGI%%@dirrm %%WWWDIR%%/nut
%%NUT_CGI%%@dirrmtry %%WWWDIR%%
