# sysutils/freenas-dispatcher/pkg-plist
%%PYTHON_SITELIBDIR%%/dispatcher_client-10.0-py2.7.egg-info
%%PYTHON_SITELIBDIR%%/dispatcherclient.py
%%PYTHON_SITELIBDIR%%/dispatcherclient.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/rpc.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/rpc.py
%%PYTHON_SITELIBDIR%%/dispatcher/client.py
%%PYTHON_SITELIBDIR%%/dispatcher/__init__.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/__init__.py
%%PYTHON_SITELIBDIR%%/dispatcher/client.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/validator.py
%%PYTHON_SITELIBDIR%%/dispatcher/validator.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/jsonenc.py
%%PYTHON_SITELIBDIR%%/dispatcher/jsonenc.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/shell.py
%%PYTHON_SITELIBDIR%%/dispatcher/shell.pyc
%%PYTHON_SITELIBDIR%%/dispatcher/file.py
%%PYTHON_SITELIBDIR%%/dispatcher/file.pyc
lib/dispatcher/src/main.py
lib/dispatcher/src/api/__init__.py
lib/dispatcher/src/api/handler.py
lib/dispatcher/src/task.py
lib/dispatcher/src/event.py
lib/dispatcher/src/cache.py
lib/dispatcher/src/query.py
lib/dispatcher/src/lib/freebsd.py
lib/dispatcher/src/lib/geom.py
lib/dispatcher/src/lib/__init__.py
lib/dispatcher/src/lib/system.py
lib/dispatcher/src/balancer.py
lib/dispatcher/src/auth.py
lib/dispatcher/src/services.py
lib/dispatcher/src/schemas.py
lib/dispatcher/src/resources.py
lib/dispatcher/src/utils.py
lib/dispatcher/src/taskproxy/main.py
lib/dispatcher/src/frontend/templates/nav.html
lib/dispatcher/src/frontend/templates/tasks.html
lib/dispatcher/src/frontend/templates/events.html
lib/dispatcher/src/frontend/templates/rpc.html
lib/dispatcher/src/frontend/templates/term.html
lib/dispatcher/src/frontend/templates/stats.html
lib/dispatcher/src/frontend/templates/apidoc/index.html
lib/dispatcher/src/frontend/templates/apidoc/nav.html
lib/dispatcher/src/frontend/templates/apidoc/rpc.html
lib/dispatcher/src/frontend/templates/apidoc/events.html
lib/dispatcher/src/frontend/templates/apidoc/tasks.html
lib/dispatcher/src/frontend/templates/apidoc/schemas.html
lib/dispatcher/src/frontend/static/middleware.js
lib/dispatcher/src/frontend/static/emitter.js
lib/dispatcher/src/frontend/static/term.js
lib/dispatcher/src/frontend/static/jquery.jsonview.js
lib/dispatcher/src/frontend/static/jquery.jsonview.css
lib/dispatcher/src/frontend/static/main.css
lib/dispatcher/src/frontend/frontend.py
lib/dispatcher/src/frontend/__init__.py
lib/dispatcher/plugins/AlertVolumePlugin.py
lib/dispatcher/plugins/SystemInfoPlugin.py
lib/dispatcher/plugins/ZfsPlugin.py
lib/dispatcher/plugins/DevdPlugin.py
lib/dispatcher/plugins/ServiceManagePlugin.py
lib/dispatcher/plugins/SSHPlugin.py
lib/dispatcher/plugins/ConfigPlugin.py
lib/dispatcher/plugins/CryptoPlugin.py
lib/dispatcher/plugins/AuditPipePlugin.py
lib/dispatcher/plugins/UserPlugin.py
lib/dispatcher/plugins/FilesystemPlugin.py
lib/dispatcher/plugins/VolumePlugin.py
lib/dispatcher/plugins/DiskPlugin.py
lib/dispatcher/plugins/NetworkPlugin.py
lib/dispatcher/plugins/EntitySubscriberPlugin.py
lib/dispatcher/plugins/UITransactionPlugin.py
lib/dispatcher/plugins/SessionPlugin.py
lib/dispatcher/plugins/ShareNFSPlugin.py
lib/dispatcher/plugins/ShareAFPPlugin.py
lib/dispatcher/plugins/SharingPlugin.py
lib/dispatcher/plugins/SystemDatasetPlugin.py
lib/dispatcher/plugins/SwapPlugin.py
lib/dispatcher/plugins/UpdatePlugin.py
lib/dispatcher/plugins/AlertPlugin.py
lib/dispatcher/plugins/MailPlugin.py
lib/dispatcher/plugins/BootPlugin.py
lib/dispatcher/plugins/NTPPlugin.py
lib/dispatcher/plugins/TunablePlugin.py
lib/dispatcher/plugins/SupportPlugin.py
lib/dispatcher/plugins/AFPPlugin.py
lib/dispatcher/plugins/CIFSPlugin.py
lib/dispatcher/plugins/CalendarTasksPlugin.py
lib/dispatcher/plugins/IPMIPlugin.py
lib/dispatcher/plugins/DynDNSPlugin.py
lib/dispatcher/plugins/FTPPlugin.py
lib/dispatcher/plugins/IPFSPlugin.py
lib/dispatcher/plugins/NFSPlugin.py
lib/dispatcher/plugins/RsyncdPlugin.py
lib/dispatcher/plugins/SMARTPlugin.py
lib/dispatcher/plugins/ShareCIFSPlugin.py
lib/dispatcher/plugins/ReplicationPlugin.py
lib/dispatcher/plugins/SNMPPlugin.py
lib/dispatcher/plugins/HAProxyPlugin.py
lib/dispatcher/plugins/LLDPPlugin.py
lib/dispatcher/plugins/RIAKCSPlugin.py
lib/dispatcher/plugins/RIAKPlugin.py
lib/dispatcher/plugins/SWIFTPlugin.py
lib/dispatcher/plugins/StanchionPlugin.py
lib/dispatcher/plugins/TFTPPlugin.py
lib/dispatcher/plugins/UPSPlugin.py
lib/dispatcher/plugins/WebDAVPlugin.py
libexec/dispatcher/notify-pam
etc/rc.d/dispatcher
sbin/dispatcherd
sbin/dispatcherctl
libexec/taskproxy
etc/middleware.conf
