# comms/soundmodem/pkg-plist
bin/soundmodemconfig
sbin/soundmodem
etc/rc.d/soundmodem
include/modem.h
include/simd.h
share/locale/fr/LC_MESSAGES/soundmodem.mo
share/locale/sv/LC_MESSAGES/soundmodem.mo
man/man8/soundmodem.8.gz
man/man8/soundmodemconfig.8.gz
@dirrmtry share/locale/fr/LC_MESSAGES
@dirrmtry share/locale/sv/LC_MESSAGES
@dirrmtry share/locale/fr
@dirrmtry share/locale/sv
@dirrmtry share/locale
