# databases/cego/pkg-plist
bin/cego
bin/cgadm
bin/cgblow
bin/cgclt
bin/cglog
etc/rc.d/cego
include/cego/Cego.h
include/cego/CegoAVLIndexCursor.h
include/cego/CegoAVLIndexEntry.h
include/cego/CegoAVLIndexManager.h
include/cego/CegoAction.h
include/cego/CegoAdm.h
include/cego/CegoAdmAction.h
include/cego/CegoAdmNet.h
include/cego/CegoAdminHandler.h
include/cego/CegoAdminThread.h
include/cego/CegoAdminThreadPool.h
include/cego/CegoAggregation.h
include/cego/CegoAlterDesc.h
include/cego/CegoAttrComp.h
include/cego/CegoAttrCond.h
include/cego/CegoAttrDesc.h
include/cego/CegoBTreeCursor.h
include/cego/CegoBTreeManager.h
include/cego/CegoBTreeNode.h
include/cego/CegoBTreeObject.h
include/cego/CegoBTreeValue.h
include/cego/CegoBeatConnection.h
include/cego/CegoBeatThread.h
include/cego/CegoBlob.h
include/cego/CegoBufferPage.h
include/cego/CegoBufferPool.h
include/cego/CegoBufferPoolEntry.h
include/cego/CegoCaseCond.h
include/cego/CegoCheckObject.h
include/cego/CegoCheckpoint.h
include/cego/CegoClob.h
include/cego/CegoComparison.h
include/cego/CegoCondDesc.h
include/cego/CegoContentObject.h
include/cego/CegoDataPointer.h
include/cego/CegoDataType.h
include/cego/CegoDatabaseFormater.h
include/cego/CegoDatabaseManager.h
include/cego/CegoDbHandler.h
include/cego/CegoDbThread.h
include/cego/CegoDbThreadPool.h
include/cego/CegoDecodableObject.h
include/cego/CegoDefs.h
include/cego/CegoDistCursor.h
include/cego/CegoDistDbHandler.h
include/cego/CegoDistManager.h
include/cego/CegoException.h
include/cego/CegoExpOutStream.h
include/cego/CegoExpr.h
include/cego/CegoFactor.h
include/cego/CegoField.h
include/cego/CegoFieldValue.h
include/cego/CegoFileHandler.h
include/cego/CegoFunction.h
include/cego/CegoGroupCursor.h
include/cego/CegoGroupNode.h
include/cego/CegoGroupSpace.h
include/cego/CegoHavingDesc.h
include/cego/CegoImpInStream.h
include/cego/CegoJDBCInfo.h
include/cego/CegoJoinObject.h
include/cego/CegoKeyObject.h
include/cego/CegoLockHandler.h
include/cego/CegoLogHandler.h
include/cego/CegoLogManager.h
include/cego/CegoLogRecord.h
include/cego/CegoLogThread.h
include/cego/CegoLogThreadPool.h
include/cego/CegoMediatorThread.h
include/cego/CegoModule.h
include/cego/CegoNet.h
include/cego/CegoObject.h
include/cego/CegoObjectCursor.h
include/cego/CegoObjectManager.h
include/cego/CegoOrderCursor.h
include/cego/CegoOrderNode.h
include/cego/CegoOrderSpace.h
include/cego/CegoOutput.h
include/cego/CegoPredDesc.h
include/cego/CegoProcAssignStmt.h
include/cego/CegoProcBlock.h
include/cego/CegoProcBlockStmt.h
include/cego/CegoProcCond.h
include/cego/CegoProcCursor.h
include/cego/CegoProcCursorCloseStmt.h
include/cego/CegoProcCursorCreateStmt.h
include/cego/CegoProcException.h
include/cego/CegoProcFetch.h
include/cego/CegoProcIfStmt.h
include/cego/CegoProcNoopStmt.h
include/cego/CegoProcObject.h
include/cego/CegoProcPred.h
include/cego/CegoProcQueryStmt.h
include/cego/CegoProcReturnStmt.h
include/cego/CegoProcStmt.h
include/cego/CegoProcThrowStmt.h
include/cego/CegoProcVar.h
include/cego/CegoProcWhileStmt.h
include/cego/CegoProcedure.h
include/cego/CegoQuery.h
include/cego/CegoQueryCache.h
include/cego/CegoQueryException.h
include/cego/CegoQueryHelper.h
include/cego/CegoRecoveryManager.h
include/cego/CegoReturnVar.h
include/cego/CegoSelect.h
include/cego/CegoSerial.h
include/cego/CegoSystemObject.h
include/cego/CegoTableCache.h
include/cego/CegoTableCursor.h
include/cego/CegoTableManager.h
include/cego/CegoTableObject.h
include/cego/CegoTerm.h
include/cego/CegoTransactionManager.h
include/cego/CegoTupleState.h
include/cego/CegoTypeConverter.h
include/cego/CegoView.h
include/cego/CegoViewObject.h
include/cego/CegoXMLHelper.h
include/cego/CegoXMLSpace.h
include/cego/CegoXMLdef.h
include/cego/CegoXPorter.h
include/cego/cgwrap.h
lib/libcego.a
lib/libcego.so
lib/libcego.so.2
lib/libcego.so.2.0.0
lib/libcgwrap.a
lib/libcgwrap.so
lib/libcgwrap.so.2
lib/libcgwrap.so.2.0.0
%%DATADIR%%/cgbackup
%%DATADIR%%/cgmkdb
%%DATADIR%%/cgrecover
%%DATADIR%%/cgsimload
%%DATADIR%%/startdb
%%DATADIR%%/stopdb
