# mail/cyrus-imapd24/pkg-plist
bin/cyradm
bin/imtest
cyrus/bin/pwcheck
cyrus/bin/arbitron
cyrus/bin/collectnews
cyrus/bin/cyrquota
cyrus/bin/deliver
cyrus/bin/dump_deliverdb
cyrus/bin/feedcyrus
cyrus/bin/fud
cyrus/bin/imapd
cyrus/bin/pop3d
cyrus/bin/reconstruct
cyrus/bin/syncnews
cyrus/bin/updateimsp
etc/imapd.conf.dist
etc/inetd.conf.cyrus
etc/rc.d/cyrus.sh
include/cyrus/acl.h
include/cyrus/assert.h
include/cyrus/auth.h
include/cyrus/bsearch.h
include/cyrus/charset.h
include/cyrus/glob.h
include/cyrus/gmtoff.h
include/cyrus/imclient.h
include/cyrus/imparse.h
include/cyrus/lock.h
include/cyrus/map.h
include/cyrus/mkgmtime.h
include/cyrus/nonblock.h
include/cyrus/parseaddr.h
include/cyrus/prot.h
include/cyrus/retry.h
include/cyrus/sasl.h
include/cyrus/sysexits.h
include/cyrus/util.h
include/cyrus/xmalloc.h
lib/libcyrus.a
@dirrm include/cyrus
@dirrm cyrus/bin
@dirrm cyrus
