# news/nntpswitch/pkg-plist
etc/nntpswitch/access.conf-dist
etc/nntpswitch/nntpswitch.conf-dist
etc/nntpswitch/servers.conf-dist
etc/rc.d/nntpswitch.sh-dist
lib/nntpswitch/acct_none.so
lib/nntpswitch/acct_syslog.so
lib/nntpswitch/auth_ignore.so
lib/nntpswitch/auth_passwd.so
lib/nntpswitch/auth_pop3.so
lib/nntpswitch/auth_remote.so
sbin/nntpswitchd
sbin/nsstats.pl
sbin/nsstats.sh
sbin/updategroups
@dirrm lib/nntpswitch
@unexec rmdir %D/etc/nntpswitch 2> /dev/null || echo "===> If you plan to do not reinstall nntpswitch, you can safely remove %D/etc/nntpswitch."
