# sysutils/samesame/pkg-plist
%%EXPERIMENTAL%%bin/fsort
bin/samearchive
bin/samearchive-lite
bin/samechflags
bin/samechmod
bin/samechown
bin/samecp
bin/samedelay
bin/samefile
bin/sameln
bin/samemv
bin/samerm
%%DISK_STORAGE%%etc/rc.d/samesame.sh
man/man1/samearchive-lite.1.gz
man/man1/samearchive.1.gz
man/man1/samechflags.1.gz
man/man1/samechmod.1.gz
man/man1/samechown.1.gz
man/man1/samecp.1.gz
man/man1/samedelay.1.gz
man/man1/samefile.1.gz
man/man1/sameln.1.gz
man/man1/samemv.1.gz
man/man1/samerm.1.gz
man/man1/samesame.1.gz
