# net-mgmt/nagios/pkg-plist
bin/netsaint
etc/netsaint/netsaint.cfg
etc/netsaint/hosts.cfg
etc/netsaint/functions
etc/rc.d/netsaint.sh
share/netsaint/docs/about.html
share/netsaint/docs/addons.html
share/netsaint/docs/cgiauth.html
share/netsaint/docs/cgis.html
share/netsaint/docs/config.html
share/netsaint/docs/configcgi.html
share/netsaint/docs/confighost.html
share/netsaint/docs/configmain.html
share/netsaint/docs/eventhandlers.html
share/netsaint/docs/extcommands.html
share/netsaint/docs/faqs.html
share/netsaint/docs/funstuff.html
share/netsaint/docs/images/cmd_cgi.gif
share/netsaint/docs/images/config_cgi.gif
share/netsaint/docs/images/extinfo_cgi_a.gif
share/netsaint/docs/images/extinfo_cgi_b.gif
share/netsaint/docs/images/extinfo_cgi_c.gif
share/netsaint/docs/images/extinfo_cgi_d.gif
share/netsaint/docs/images/history_cgi.gif
share/netsaint/docs/images/logic-host-event-handler.gif
share/netsaint/docs/images/logic-host-notification.gif
share/netsaint/docs/images/logic-host-state-handler.gif
share/netsaint/docs/images/logic-route-check.gif
share/netsaint/docs/images/logic-service-event-handler.gif
share/netsaint/docs/images/logic-service-notification.gif
share/netsaint/docs/images/logic-service-problem.gif
share/netsaint/docs/images/logic-service-recovery.gif
share/netsaint/docs/images/logic-top.gif
share/netsaint/docs/images/network-heirarchy.gif
share/netsaint/docs/images/notification-process.gif
share/netsaint/docs/images/notifications_cgi.gif
share/netsaint/docs/images/physical-network.gif
share/netsaint/docs/images/redundancy-1.gif
share/netsaint/docs/images/redundancy-2.gif
share/netsaint/docs/images/relationship-complex.gif
share/netsaint/docs/images/relationship-simple.gif
share/netsaint/docs/images/showlog_cgi.gif
share/netsaint/docs/images/status_cgi_a.gif
share/netsaint/docs/images/status_cgi_b.gif
share/netsaint/docs/images/statusmap_cgi.gif
share/netsaint/docs/images/statuswrl_cgi.gif
share/netsaint/docs/index.html
share/netsaint/docs/installing.html
share/netsaint/docs/installweb.html
share/netsaint/docs/macros.html
share/netsaint/docs/makehtml.html
share/netsaint/docs/networkreachability.html
share/netsaint/docs/notifications.html
share/netsaint/docs/parallelization.html
share/netsaint/docs/pluginhowto.html
share/netsaint/docs/plugins.html
share/netsaint/docs/programmodes.html
share/netsaint/docs/redundancy.html
share/netsaint/docs/robots.txt
share/netsaint/docs/starting.html
share/netsaint/docs/statetypes.html
share/netsaint/docs/statuslevels.html
share/netsaint/docs/stoprestart.html
share/netsaint/docs/theory-network.html
share/netsaint/docs/theory-time.html
share/netsaint/docs/theory.html
share/netsaint/docs/timeperiods.html
share/netsaint/docs/verifyconfig.html
share/netsaint/docs/whatsnew.html
share/netsaint/images/active.gif
share/netsaint/images/command.gif
share/netsaint/images/commandsbar.gif
share/netsaint/images/comment.gif
share/netsaint/images/config.gif
share/netsaint/images/contact.gif
share/netsaint/images/contactgroup.gif
share/netsaint/images/contactgroupsbar.gif
share/netsaint/images/contactsbar.gif
share/netsaint/images/critical.gif
share/netsaint/images/delay.gif
share/netsaint/images/delete.gif
share/netsaint/images/disabled.gif
share/netsaint/images/docbar.gif
share/netsaint/images/enabled.gif
share/netsaint/images/greendot.gif
share/netsaint/images/history.gif
share/netsaint/images/historybar.gif
share/netsaint/images/homebar.gif
share/netsaint/images/host.gif
share/netsaint/images/hostevent.gif
share/netsaint/images/hostgroupsbar.gif
share/netsaint/images/hostsbar.gif
share/netsaint/images/info.gif
share/netsaint/images/left.gif
share/netsaint/images/logfilebar.gif
share/netsaint/images/logofullsize.gif
share/netsaint/images/logos/aix.gd2
share/netsaint/images/logos/aix.gif
share/netsaint/images/logos/amiga.gd2
share/netsaint/images/logos/amiga.gif
share/netsaint/images/logos/apple.gd2
share/netsaint/images/logos/apple.gif
share/netsaint/images/logos/beos.gd2
share/netsaint/images/logos/beos.gif
share/netsaint/images/logos/caldera.gd2
share/netsaint/images/logos/caldera.gif
share/netsaint/images/logos/cat1900.gd2
share/netsaint/images/logos/cat2900.gd2
share/netsaint/images/logos/cat5000.gd2
share/netsaint/images/logos/debian.gd2
share/netsaint/images/logos/debian.gif
share/netsaint/images/logos/freebsd40.gd2
share/netsaint/images/logos/freebsd40.gif
share/netsaint/images/logos/hp-printer40.gd2
share/netsaint/images/logos/hp-printer40.gif
share/netsaint/images/logos/hpux.gd2
share/netsaint/images/logos/hpux.gif
share/netsaint/images/logos/irix.gd2
share/netsaint/images/logos/irix.gif
share/netsaint/images/logos/linux40.gd2
share/netsaint/images/logos/linux40.gif
share/netsaint/images/logos/logo.gd2
share/netsaint/images/logos/mac40.gd2
share/netsaint/images/logos/mac40.gif
share/netsaint/images/logos/mandrake.gd2
share/netsaint/images/logos/mandrake.gif
share/netsaint/images/logos/next.gd2
share/netsaint/images/logos/next.gif
share/netsaint/images/logos/ng-switch40.gd2
share/netsaint/images/logos/ng-switch40.gif
share/netsaint/images/logos/novell40.gd2
share/netsaint/images/logos/novell40.gif
share/netsaint/images/logos/openbsd.gd2
share/netsaint/images/logos/openbsd.gif
share/netsaint/images/logos/redhat.gd2
share/netsaint/images/logos/redhat.gif
share/netsaint/images/logos/router40.gd2
share/netsaint/images/logos/router40.gif
share/netsaint/images/logos/slackware.gd2
share/netsaint/images/logos/slackware.gif
share/netsaint/images/logos/stampede.gd2
share/netsaint/images/logos/stampede.gif
share/netsaint/images/logos/storm.gd2
share/netsaint/images/logos/storm.gif
share/netsaint/images/logos/sun40.gd2
share/netsaint/images/logos/sun40.gif
share/netsaint/images/logos/sunlogo.gd2
share/netsaint/images/logos/sunlogo.gif
share/netsaint/images/logos/switch40.gd2
share/netsaint/images/logos/switch40.gif
share/netsaint/images/logos/turbolinux.gd2
share/netsaint/images/logos/turbolinux.gif
share/netsaint/images/logos/ultrapenguin.gd2
share/netsaint/images/logos/ultrapenguin.gif
share/netsaint/images/logos/unicos.gd2
share/netsaint/images/logos/unicos.gif
share/netsaint/images/logos/unknown.gd2
share/netsaint/images/logos/win40.gd2
share/netsaint/images/logos/win40.gif
share/netsaint/images/logos/yellowdog.gd2
share/netsaint/images/logos/yellowdog.gif
share/netsaint/images/logrotate.gif
share/netsaint/images/logs.gif
share/netsaint/images/ndisabled.gif
share/netsaint/images/networkhealthbar.gif
share/netsaint/images/notes.gif
share/netsaint/images/notificationsbar.gif
share/netsaint/images/notify.gif
share/netsaint/images/orangedot.gif
share/netsaint/images/parentup.gif
share/netsaint/images/processinfobar.gif
share/netsaint/images/recovery.gif
share/netsaint/images/restart.gif
share/netsaint/images/right.gif
share/netsaint/images/sbbg.gif
share/netsaint/images/sbconfig.gif
share/netsaint/images/sbgeneral.gif
share/netsaint/images/sblogo.gif
share/netsaint/images/sbmonitor.gif
share/netsaint/images/scale.gif
share/netsaint/images/serviceevent.gif
share/netsaint/images/servicesbar.gif
share/netsaint/images/standby.gif
share/netsaint/images/start.gif
share/netsaint/images/statudetailbar.gif
share/netsaint/images/status.gif
share/netsaint/images/status2.gif
share/netsaint/images/status3.gif
share/netsaint/images/status4.gif
share/netsaint/images/statusdetailbar.gif
share/netsaint/images/statusmapbar.gif
share/netsaint/images/statusoverviewbar.gif
share/netsaint/images/statusworldbar.gif
share/netsaint/images/stop.gif
share/netsaint/images/thermcrit.gif
share/netsaint/images/thermok.gif
share/netsaint/images/thermwarn.gif
share/netsaint/images/timeperiodsbar.gif
share/netsaint/images/unknown.gif
share/netsaint/images/warning.gif
share/netsaint/images/xinfo.gif
share/netsaint/images/zombie.gif
share/netsaint/index.html
share/netsaint/main.html
share/netsaint/robots.txt
share/netsaint/sbin/cmd.cgi
share/netsaint/sbin/config.cgi
share/netsaint/sbin/extinfo.cgi
share/netsaint/sbin/history.cgi
share/netsaint/sbin/notifications.cgi
share/netsaint/sbin/nscgi.cfg
share/netsaint/sbin/showlog.cgi
share/netsaint/sbin/status.cgi
share/netsaint/sbin/statusmap.cgi
share/netsaint/sbin/statuswrl.cgi
share/netsaint/side.html
@dirrm share/netsaint/sbin
@dirrm share/netsaint/images/logos
@dirrm share/netsaint/images
@dirrm share/netsaint/docs/images
@dirrm share/netsaint/docs
@dirrm share/netsaint
