# net-mgmt/nagios/pkg-plist
@exec mkdir -p %D/share/nagios/ssi
@exec mkdir -p /var/spool/nagios/archives
@exec mkdir -p /var/spool/nagios/rw
@exec chmod 775 /var/spool/nagios /var/spool/nagios/archives /var/spool/nagios/rw
@exec chown nagios:nagios /var/spool/nagios /var/spool/nagios/archives
@exec chown nagios:www /var/spool/nagios/rw
bin/nagios
bin/nagiostats
etc/nagios/bigger.cfg-sample
etc/nagios/cgi.cfg-sample
etc/nagios/checkcommands.cfg-sample
etc/nagios/minimal.cfg-sample
etc/nagios/misccommands.cfg-sample
etc/nagios/nagios.cfg-sample
etc/nagios/resource.cfg-sample
etc/rc.d/nagios.sh
share/nagios/cgi-bin/avail.cgi
share/nagios/cgi-bin/cmd.cgi
share/nagios/cgi-bin/config.cgi
share/nagios/cgi-bin/extinfo.cgi
share/nagios/cgi-bin/histogram.cgi
share/nagios/cgi-bin/history.cgi
share/nagios/cgi-bin/notifications.cgi
share/nagios/cgi-bin/outages.cgi
share/nagios/cgi-bin/showlog.cgi
share/nagios/cgi-bin/status.cgi
share/nagios/cgi-bin/statusmap.cgi
share/nagios/cgi-bin/statuswml.cgi
share/nagios/cgi-bin/statuswrl.cgi
share/nagios/cgi-bin/summary.cgi
share/nagios/cgi-bin/tac.cgi
share/nagios/cgi-bin/trends.cgi
share/nagios/contexthelp/A1.html
share/nagios/contexthelp/A2.html
share/nagios/contexthelp/A3.html
share/nagios/contexthelp/A4.html
share/nagios/contexthelp/A5.html
share/nagios/contexthelp/A6.html
share/nagios/contexthelp/A7.html
share/nagios/contexthelp/B1.html
share/nagios/contexthelp/C1.html
share/nagios/contexthelp/D1.html
share/nagios/contexthelp/E1.html
share/nagios/contexthelp/F1.html
share/nagios/contexthelp/G1.html
share/nagios/contexthelp/G2.html
share/nagios/contexthelp/G3.html
share/nagios/contexthelp/G4.html
share/nagios/contexthelp/G5.html
share/nagios/contexthelp/G6.html
share/nagios/contexthelp/H1.html
share/nagios/contexthelp/H2.html
share/nagios/contexthelp/H3.html
share/nagios/contexthelp/H4.html
share/nagios/contexthelp/H5.html
share/nagios/contexthelp/H6.html
share/nagios/contexthelp/H7.html
share/nagios/contexthelp/H8.html
share/nagios/contexthelp/I1.html
share/nagios/contexthelp/I2.html
share/nagios/contexthelp/I3.html
share/nagios/contexthelp/I4.html
share/nagios/contexthelp/I5.html
share/nagios/contexthelp/I6.html
share/nagios/contexthelp/I7.html
share/nagios/contexthelp/I8.html
share/nagios/contexthelp/I9.html
share/nagios/contexthelp/J1.html
share/nagios/contexthelp/K1.html
share/nagios/contexthelp/L1.html
share/nagios/contexthelp/L10.html
share/nagios/contexthelp/L11.html
share/nagios/contexthelp/L12.html
share/nagios/contexthelp/L13.html
share/nagios/contexthelp/L2.html
share/nagios/contexthelp/L3.html
share/nagios/contexthelp/L4.html
share/nagios/contexthelp/L5.html
share/nagios/contexthelp/L6.html
share/nagios/contexthelp/L7.html
share/nagios/contexthelp/L8.html
share/nagios/contexthelp/L9.html
share/nagios/contexthelp/M1.html
share/nagios/contexthelp/M2.html
share/nagios/contexthelp/M3.html
share/nagios/contexthelp/M4.html
share/nagios/contexthelp/M5.html
share/nagios/contexthelp/M6.html
share/nagios/contexthelp/N1.html
share/nagios/contexthelp/N2.html
share/nagios/contexthelp/N3.html
share/nagios/contexthelp/N4.html
share/nagios/contexthelp/N5.html
share/nagios/contexthelp/N6.html
share/nagios/contexthelp/N7.html
share/nagios/docs/images/activepassive.png
share/nagios/docs/images/cgi-avail-a.png
share/nagios/docs/images/cgi-avail-b.png
share/nagios/docs/images/cgi-cmd.png
share/nagios/docs/images/cgi-config.png
share/nagios/docs/images/cgi-extinfo-a.png
share/nagios/docs/images/cgi-extinfo-b.png
share/nagios/docs/images/cgi-extinfo-c.png
share/nagios/docs/images/cgi-extinfo-d.png
share/nagios/docs/images/cgi-histogram.png
share/nagios/docs/images/cgi-history.png
share/nagios/docs/images/cgi-notifications.png
share/nagios/docs/images/cgi-outages.png
share/nagios/docs/images/cgi-showlog.png
share/nagios/docs/images/cgi-status-a.png
share/nagios/docs/images/cgi-status-b.png
share/nagios/docs/images/cgi-status-c.png
share/nagios/docs/images/cgi-status-d.png
share/nagios/docs/images/cgi-statusmap.png
share/nagios/docs/images/cgi-statuswml.png
share/nagios/docs/images/cgi-statuswrl.png
share/nagios/docs/images/cgi-summary.png
share/nagios/docs/images/cgi-tac.png
share/nagios/docs/images/cgi-trends.png
share/nagios/docs/images/checktiming.png
share/nagios/docs/images/distributed.png
share/nagios/docs/images/host-dependencies.png
share/nagios/docs/images/indirecthostcheck.png
share/nagios/docs/images/indirectsvccheck.png
share/nagios/docs/images/indirectsvccheck2.png
share/nagios/docs/images/interleaved1.png
share/nagios/docs/images/interleaved2.png
share/nagios/docs/images/interleaved3.png
share/nagios/docs/images/logofullsize.png
share/nagios/docs/images/mrtgstats.png
share/nagios/docs/images/network-heirarchy.png
share/nagios/docs/images/network-outage1.png
share/nagios/docs/images/network-outage2.png
share/nagios/docs/images/noninterleaved1.png
share/nagios/docs/images/noninterleaved2.png
share/nagios/docs/images/physical-network.png
share/nagios/docs/images/plugintheory.png
share/nagios/docs/images/redudancy.png
share/nagios/docs/images/redundancy.png
share/nagios/docs/images/service-dependencies.png
share/nagios/docs/images/statetransitions.png
share/nagios/docs/images/statetransitions2.png
share/nagios/docs/about.html
share/nagios/docs/addons.html
share/nagios/docs/adaptive.html
share/nagios/docs/beginners.html
share/nagios/docs/cgiauth.html
share/nagios/docs/cgiincludes.html
share/nagios/docs/cgis.html
share/nagios/docs/checkscheduling.html
share/nagios/docs/clusters.html
share/nagios/docs/commandfile.html
share/nagios/docs/config.html
share/nagios/docs/configcgi.html
share/nagios/docs/configextinfo.html
share/nagios/docs/configmain.html
share/nagios/docs/configobject.html
share/nagios/docs/dependencies.html
share/nagios/docs/distributed.html
share/nagios/docs/downtime.html
share/nagios/docs/embeddedperl.html
share/nagios/docs/escalations.html
share/nagios/docs/eventhandlers.html
share/nagios/docs/extcommands.html
share/nagios/docs/flapping.html
share/nagios/docs/freshness.html
share/nagios/docs/funstuff.html
share/nagios/docs/index.html
share/nagios/docs/indirectchecks.html
share/nagios/docs/installing.html
share/nagios/docs/installweb.html
share/nagios/docs/int-snmptrap.html
share/nagios/docs/int-tcpwrappers.html
share/nagios/docs/macros.html
share/nagios/docs/nagiostats.html
share/nagios/docs/networkoutages.html
share/nagios/docs/networkreachability.html
share/nagios/docs/notifications.html
share/nagios/docs/parallelization.html
share/nagios/docs/passivechecks.html
share/nagios/docs/perfdata.html
share/nagios/docs/plugins.html
share/nagios/docs/plugintheory.html
share/nagios/docs/redundancy.html
share/nagios/docs/robots.txt
share/nagios/docs/security.html
share/nagios/docs/stalking.html
share/nagios/docs/starting.html
share/nagios/docs/statetypes.html
share/nagios/docs/stoprestart.html
share/nagios/docs/templaterecursion.html
share/nagios/docs/templatetricks.html
share/nagios/docs/timeperiods.html
share/nagios/docs/toc.html
share/nagios/docs/tuning.html
share/nagios/docs/verifyconfig.html
share/nagios/docs/volatileservices.html
share/nagios/docs/whatsnew.html
share/nagios/docs/xodtemplate.html
share/nagios/images/ack.gif
share/nagios/images/action.gif
share/nagios/images/command.png
share/nagios/images/comment.gif
share/nagios/images/contexthelp1.gif
share/nagios/images/contexthelp2.gif
share/nagios/images/critical.png
share/nagios/images/delay.gif
share/nagios/images/delete.gif
share/nagios/images/detail.gif
share/nagios/images/disabled.gif
share/nagios/images/down.gif
share/nagios/images/downtime.gif
share/nagios/images/empty.gif
share/nagios/images/enabled.gif
share/nagios/images/extinfo.gif
share/nagios/images/flapping.gif
share/nagios/images/greendot.gif
share/nagios/images/histogram.png
share/nagios/images/history.gif
share/nagios/images/hostevent.gif
share/nagios/images/info.png
share/nagios/images/left.gif
share/nagios/images/logofullsize.jpg
share/nagios/images/logos/nagios.gd2
share/nagios/images/logos/nagios.gif
share/nagios/images/logos/nagiosvrml.png
share/nagios/images/logos/unknown.gd2
share/nagios/images/logos/unknown.gif
share/nagios/images/logrotate.png
share/nagios/images/ndisabled.gif
share/nagios/images/noack.gif
share/nagios/images/notes.gif
share/nagios/images/notify.gif
share/nagios/images/orangedot.gif
share/nagios/images/passiveonly.gif
share/nagios/images/recovery.png
share/nagios/images/redudancy.png
share/nagios/images/redundancy.png
share/nagios/images/restart.gif
share/nagios/images/right.gif
share/nagios/images/sblogo.jpg
share/nagios/images/serviceevent.gif
share/nagios/images/start.gif
share/nagios/images/status.gif
share/nagios/images/status2.gif
share/nagios/images/status3.gif
share/nagios/images/status4.gif
share/nagios/images/stop.gif
share/nagios/images/tacdisabled.jpg
share/nagios/images/tacdisabled.png
share/nagios/images/tacenabled.jpg
share/nagios/images/tacenabled.png
share/nagios/images/thermcrit.png
share/nagios/images/thermok.png
share/nagios/images/thermwarn.png
share/nagios/images/trends.gif
share/nagios/images/trendshost.png
share/nagios/images/trendssvc.png
share/nagios/images/unknown.png
share/nagios/images/up.gif
share/nagios/images/warning.png
share/nagios/images/weblogo1.png
share/nagios/images/zoom1.gif
share/nagios/images/zoom2.gif
share/nagios/index.html
share/nagios/main.html
share/nagios/media/critical.wav
share/nagios/media/hostdown.wav
share/nagios/media/warning.wav
share/nagios/robots.txt
share/nagios/side.html
share/nagios/stylesheets/avail.css
share/nagios/stylesheets/checksanity.css
share/nagios/stylesheets/cmd.css
share/nagios/stylesheets/common.css
share/nagios/stylesheets/config.css
share/nagios/stylesheets/extinfo.css
share/nagios/stylesheets/histogram.css
share/nagios/stylesheets/history.css
share/nagios/stylesheets/ministatus.css
share/nagios/stylesheets/notifications.css
share/nagios/stylesheets/outages.css
share/nagios/stylesheets/showlog.css
share/nagios/stylesheets/status.css
share/nagios/stylesheets/statusmap.css
share/nagios/stylesheets/summary.css
share/nagios/stylesheets/tac.css
share/nagios/stylesheets/trends.css
@unexec rmdir %D/etc/nagios 2>/dev/null || true
@dirrm share/nagios/cgi-bin
@dirrm share/nagios/contexthelp
@dirrm share/nagios/docs/images
@dirrm share/nagios/docs
@dirrm share/nagios/images/logos
@dirrm share/nagios/images
@dirrm share/nagios/media
@dirrm share/nagios/ssi
@dirrm share/nagios/stylesheets
@unexec rmdir %D/share/nagios 2>/dev/null || true
@unexec rmdir /var/spool/nagios/archives 2>/dev/null || true
@unexec rmdir /var/spool/nagios/rw 2>/dev/null || true
@unexec rmdir /var/spool/nagios 2>/dev/null || true
