# lang/rexx-regina/files/openrc-rxstack.in
#!/sbin/openrc-run
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# rxstack_flags (str):   Set to "-d" by default.

name=rxstack
description="Regina REXX queue daemon"

command=%%PREFIX%%/bin/${name}
command_args=-d

depend() {
	keyword -shutdown
}
