# lang/rexx-regina/files/rxstack.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: rxstack
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# rxstack_enable (bool):   Set to NO by default.
#               Set it to YES to enable rxstack.
# rxstack_flags (str):   Set to "-d" by default.
#

. /etc/rc.subr

name=rxstack
rcvar=rxstack_enable
desc="Regina REXX queue daemon"

load_rc_config $name

: ${rxstack_enable:=NO}

command=%%PREFIX%%/bin/${name}
command_args=-d

run_rc_command "$1"
