# www/fhttpd/pkg-plist
etc/fhttpd.conf
etc/rc.d/fhttpd.sh
%%LOGFILE%%
sbin/fhttpd
sbin/fhttpd-file
sbin/fhttpd-log
%%DATADIR%%/bin/finger-nohtml
%%DATADIR%%/cgi-bin/.index.html
%%DATADIR%%/cgi-bin/finger
%%DATADIR%%/cgi-bin/fortune
%%DATADIR%%/icons/at_work.gif
%%DATADIR%%/icons/audio.gif
%%DATADIR%%/icons/audio.xbm
%%DATADIR%%/icons/binary.gif
%%DATADIR%%/icons/binary.xbm
%%DATADIR%%/icons/c_text.gif
%%DATADIR%%/icons/c_text.xbm
%%DATADIR%%/icons/cdromicon.gif
%%DATADIR%%/icons/exec.gif
%%DATADIR%%/icons/folder.gif
%%DATADIR%%/icons/folder.xbm
%%DATADIR%%/icons/ftp.gif
%%DATADIR%%/icons/html.gif
%%DATADIR%%/icons/html.xbm
%%DATADIR%%/icons/image.gif
%%DATADIR%%/icons/image.xbm
%%DATADIR%%/icons/movie.gif
%%DATADIR%%/icons/text.gif
%%DATADIR%%/icons/text.xbm
%%DATADIR%%/index.html
%%DATADIR%%/noaccess.html
%%DATADIR%%/private/.footer.html
%%DATADIR%%/private/.header.html
%%DATADIR%%/private/.index.html
%%DATADIR%%/private/dummy_file
%%DATADIR%%/private/dummy_file.html
%%DATADIR%%/pub/.footer.html
%%DATADIR%%/pub/.header.html
%%DATADIR%%/pub/.index.html
%%DATADIR%%/pub/dummy_file.html
@dirrm %%DATADIR%%/bin
@dirrm %%DATADIR%%/cgi-bin
@dirrm %%DATADIR%%/icons
@dirrm %%DATADIR%%/incoming
@dirrm %%DATADIR%%/private
@dirrm %%DATADIR%%/pub
@dirrm %%DATADIR%%
