# sysutils/qjail/pkg-plist
bin/qjail
bin/qjail.vnet.be
bin/qjail.vnet.ng
etc/qjail.portsnap.conf
etc/rc.d/qjail.bootime
%%EXAMPLESDIR%%/upgrade-info.txt
%%EXAMPLESDIR%%/vnet/jail.amd64
%%EXAMPLESDIR%%/vnet/jail.i386
%%EXAMPLESDIR%%/vnet/ipfw.rules.host
%%EXAMPLESDIR%%/vnet/ipfw.rules.vnet
%%EXAMPLESDIR%%/vnet/ipf.rules.host
%%EXAMPLESDIR%%/vnet/ipf.rules.vnet
%%EXAMPLESDIR%%/vnet/pf.rules.host
%%EXAMPLESDIR%%/vnet/pf.rules.vnet
%%EXAMPLESDIR%%/default/etc/make.conf
%%EXAMPLESDIR%%/default/etc/motd
%%EXAMPLESDIR%%/default/etc/periodic.conf
%%EXAMPLESDIR%%/default/etc/rc.conf
%%EXAMPLESDIR%%/default/root/.cshrc
%%EXAMPLESDIR%%/default/usr/local/etc/sudoers
%%EXAMPLESDIR%%/ssh-default/etc/group
%%EXAMPLESDIR%%/ssh-default/etc/make.conf
%%EXAMPLESDIR%%/ssh-default/etc/master.passwd
%%EXAMPLESDIR%%/ssh-default/etc/motd
%%EXAMPLESDIR%%/ssh-default/etc/passwd
%%EXAMPLESDIR%%/ssh-default/etc/periodic.conf
%%EXAMPLESDIR%%/ssh-default/etc/pwd.db
%%EXAMPLESDIR%%/ssh-default/etc/rc.conf
%%EXAMPLESDIR%%/ssh-default/etc/spwd.db
%%EXAMPLESDIR%%/ssh-default/etc/ssh/sshd_config
%%EXAMPLESDIR%%/ssh-default/root/.cshrc
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.cshrc
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.login
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.login_conf
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.mail_aliases
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.mailrc
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.profile
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.rhosts
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.shrc
%%EXAMPLESDIR%%/ssh-default/usr/local/etc/sudoers
@dirrm %%EXAMPLESDIR%%/default/usr/local/etc
@dirrm %%EXAMPLESDIR%%/default/usr/local
@dirrm %%EXAMPLESDIR%%/default/usr
@dirrm %%EXAMPLESDIR%%/default/root
@dirrm %%EXAMPLESDIR%%/default/etc
@dirrm %%EXAMPLESDIR%%/default
@dirrm %%EXAMPLESDIR%%/ssh-default/usr/home/qjail
@dirrm %%EXAMPLESDIR%%/ssh-default/usr/home
@dirrm %%EXAMPLESDIR%%/ssh-default/usr/local/etc
@dirrm %%EXAMPLESDIR%%/ssh-default/usr/local
@dirrm %%EXAMPLESDIR%%/ssh-default/usr
@dirrm %%EXAMPLESDIR%%/ssh-default/root
@dirrm %%EXAMPLESDIR%%/ssh-default/etc/ssh
@dirrm %%EXAMPLESDIR%%/ssh-default/etc
@dirrm %%EXAMPLESDIR%%/ssh-default
@dirrm %%EXAMPLESDIR%%/vnet
@dirrm %%EXAMPLESDIR%%
