# sysutils/qjail/pkg-plist
bin/qjail
bin/qjail.vnet.be
bin/qjail.vnet.ng
@sample etc/qjail.portsnap.conf.sample
etc/rc.d/qjail.bootime
man/man8/qjail.8.gz
man/man8/qjail-intro.8.gz
man/man8/qjail-howto.8.gz
%%EXAMPLESDIR%%/vnet/ipfw.rules.host
%%EXAMPLESDIR%%/vnet/ipfw.rules.vnet
%%EXAMPLESDIR%%/vnet/ipf.rules.host
%%EXAMPLESDIR%%/vnet/ipf.rules.vnet
%%EXAMPLESDIR%%/vnet/pf.rules.host
%%EXAMPLESDIR%%/vnet/pf.rules.vnet
%%EXAMPLESDIR%%/default/etc/make.conf
%%EXAMPLESDIR%%/default/etc/motd
%%EXAMPLESDIR%%/default/etc/periodic.conf
%%EXAMPLESDIR%%/default/etc/periodic103.conf
%%EXAMPLESDIR%%/default/etc/rc.conf
%%EXAMPLESDIR%%/default/etc/rc103.conf
%%EXAMPLESDIR%%/default/root/.cshrc
%%EXAMPLESDIR%%/default/usr/local/etc/pkg.conf
%%EXAMPLESDIR%%/default/usr/local/etc/sudoers
%%EXAMPLESDIR%%/ssh-default/etc/group
%%EXAMPLESDIR%%/ssh-default/etc/make.conf
%%EXAMPLESDIR%%/ssh-default/etc/master.passwd
%%EXAMPLESDIR%%/ssh-default/etc/motd
%%EXAMPLESDIR%%/ssh-default/etc/passwd
%%EXAMPLESDIR%%/ssh-default/etc/periodic.conf
%%EXAMPLESDIR%%/ssh-default/etc/periodic103.conf
%%EXAMPLESDIR%%/ssh-default/etc/pwd.db
%%EXAMPLESDIR%%/ssh-default/etc/rc.conf
%%EXAMPLESDIR%%/ssh-default/etc/rc103.conf
%%EXAMPLESDIR%%/ssh-default/etc/spwd.db
%%EXAMPLESDIR%%/ssh-default/etc/ssh/sshd_config
%%EXAMPLESDIR%%/ssh-default/root/.cshrc
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.cshrc
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.login
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.login_conf
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.mail_aliases
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.mailrc
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.profile
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.rhosts
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.shrc
%%EXAMPLESDIR%%/ssh-default/usr/local/etc/pkg.conf
%%EXAMPLESDIR%%/ssh-default/usr/local/etc/sudoers
