# sysutils/qjail/pkg-plist
bin/qjail
etc/qjail.conf.sample
etc/qjail.portsnap.conf
etc/rc.d/qjail-jail2
etc/rc.d/qjail2
%%EXAMPLESDIR%%/default/etc/make.conf
%%EXAMPLESDIR%%/default/etc/motd
%%EXAMPLESDIR%%/default/etc/periodic.conf
%%EXAMPLESDIR%%/default/etc/rc.conf
%%EXAMPLESDIR%%/default/root/.cshrc
%%EXAMPLESDIR%%/default/usr/local/etc/sudoers
%%EXAMPLESDIR%%/nullmailer-example/etc/mail/mailer.conf
%%EXAMPLESDIR%%/nullmailer-example/etc/rc.conf
%%EXAMPLESDIR%%/nullmailer-example/usr/local/etc/nullmailer/remotes
%%EXAMPLESDIR%%/ssh-default/etc/group
%%EXAMPLESDIR%%/ssh-default/etc/make.conf
%%EXAMPLESDIR%%/ssh-default/etc/master.passwd
%%EXAMPLESDIR%%/ssh-default/etc/motd
%%EXAMPLESDIR%%/ssh-default/etc/passwd
%%EXAMPLESDIR%%/ssh-default/etc/periodic.conf
%%EXAMPLESDIR%%/ssh-default/etc/pwd.db
%%EXAMPLESDIR%%/ssh-default/etc/rc.conf
%%EXAMPLESDIR%%/ssh-default/etc/spwd.db
%%EXAMPLESDIR%%/ssh-default/etc/ssh/sshd_config
%%EXAMPLESDIR%%/ssh-default/root/.cshrc
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.cshrc
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.login
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.login_conf
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.mail_aliases
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.mailrc
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.profile
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.rhosts
%%EXAMPLESDIR%%/ssh-default/usr/home/qjail/.shrc
%%EXAMPLESDIR%%/ssh-default/usr/local/etc/sudoers
@dirrm %%EXAMPLESDIR%%/ssh-default/usr/home/qjail
@dirrm %%EXAMPLESDIR%%/ssh-default/usr/home
@dirrm %%EXAMPLESDIR%%/ssh-default/usr/local/etc
@dirrm %%EXAMPLESDIR%%/ssh-default/usr/local
@dirrm %%EXAMPLESDIR%%/ssh-default/usr
@dirrm %%EXAMPLESDIR%%/ssh-default/root
@dirrm %%EXAMPLESDIR%%/ssh-default/etc/ssh
@dirrm %%EXAMPLESDIR%%/ssh-default/etc
@dirrm %%EXAMPLESDIR%%/ssh-default
@dirrm %%EXAMPLESDIR%%/nullmailer-example/usr/local/etc/nullmailer
@dirrm %%EXAMPLESDIR%%/nullmailer-example/usr/local/etc
@dirrm %%EXAMPLESDIR%%/nullmailer-example/usr/local
@dirrm %%EXAMPLESDIR%%/nullmailer-example/usr
@dirrm %%EXAMPLESDIR%%/nullmailer-example/etc/mail
@dirrm %%EXAMPLESDIR%%/nullmailer-example/etc
@dirrm %%EXAMPLESDIR%%/nullmailer-example
@dirrm %%EXAMPLESDIR%%/default/usr/local/etc
@dirrm %%EXAMPLESDIR%%/default/usr/local
@dirrm %%EXAMPLESDIR%%/default/usr
@dirrm %%EXAMPLESDIR%%/default/root
@dirrm %%EXAMPLESDIR%%/default/etc
@dirrm %%EXAMPLESDIR%%/default
@dirrm %%EXAMPLESDIR%%
