# mail/couriergraph/pkg-plist
etc/rc.d/couriergraph.sh
www/cgi-bin/couriergraph.cgi
sbin/couriergraph.pl
@unexec rmdir %D/www/cgi-bin 2>/dev/null || true
