# sysutils/bareos-server/files/bareos-fd.in
#!/bin/sh
#
# $FreeBSD: head/sysutils/bacula5-server/files/bacula-fd.in 323275 2013-07-19 09:44:58Z rm $
#
# PROVIDE: bareos_fd
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# bareos_fd_enable  (bool):  Set to NO by default.
#               Set it to YES to enable bareos_fd.
# bareos_fd_flags (params):  Set params used to start bareos_fd.
# bareos_fd_config (params):  Path to the config file/directory
#

. /etc/rc.subr

name="bareos_fd"
rcvar=${name}_enable

load_rc_config $name

: ${bareos_fd_enable="NO"}
: ${bareos_fd_flags="-u root -g wheel -v"}
: ${bareos_fd_config="%%PREFIX%%/etc/bareos/bareos-fd.d"}
: ${bareos_fd_pidfile="/var/run/bareos/bareos-fd.9102.pid"}

command=%%PREFIX%%/sbin/bareos-fd
command_args="-c ${bareos_fd_config}"
pidfile="${bareos_fd_pidfile}"

start_precmd() {
	[ -d "${pidfile%/*}" ] || {
		getent passwd bareos > /dev/null 2&>1
		if [ $? -eq 0] ; then
			# user bareos exists
			install -d -o bareos -g bareos ${pidfile%/*}
		else
			# bareos_fd as solitair service
			install -d -o root -g wheel ${pidfile%/*}
		fi
	}
}

run_rc_command "$1"
