# sysutils/bareos-server/files/bareos-dir.in
#!/bin/sh
#
# $FreeBSD: head/sysutils/bareos-server/files/bareos-dir.in 323275 2013-07-19 09:44:58Z rm $
#
# PROVIDE: bareos_dir
# REQUIRE: DAEMON %%REQ_MYSQL%% %%REQ_PGSQL%%
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# bareos_dir_enable  (bool):   Set to NO by default.
#                Set it to YES to enable bareos_dir.
# bareos_dir_flags (params):   Set params used to start bareos_dir.
# bareos_dir_config (params):  Path to the config file/directory
#

. /etc/rc.subr

name="bareos_dir"
rcvar=${name}_enable

load_rc_config $name

: ${bareos_dir_enable="NO"}
: ${bareos_dir_flags="-u bareos -g bareos -v"}
: ${bareos_dir_config="%%PREFIX%%/etc/bareos/bareos-dir.d"}
: ${bareos_dir_pidfile="/var/run/bareos/bareos-dir.9101.pid"}

command=%%PREFIX%%/sbin/bareos-dir
command_args="-c ${bareos_dir_config}"
start_precmd="bareos_dir_testconfig"
restart_precmd="bareos_dir_testconfig"
pidfile="${bareos_dir_pidfile}"

bareos_dir_testconfig()
{
	echo Performing sanity check on bareos-dir configuration:
	eval ${command} ${bareos_dir_flags} -t
}

start_precmd() {
	[ -d "${pidfile%/*}" ] || install -d -o bareos -g bareos ${pidfile%/*}
}

run_rc_command "$1"
