# mail/drac/files/dracd.sh.in
#!/bin/sh
# $FreeBSD$
#
# PROVIDE: dracd
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable postfix:
# dracd_enable (bool):		Set it to "YES" to enable rpc.dracd.
#				Default is "NO".
# dracd_flags (str):		Set to "" by default.
#                               Extra flags passed to rpc.dracd.
#

. /etc/rc.subr

name="dracd"
rcvar=dracd_enable
command="%%PREFIX%%/sbin/rpc.dracd"

load_rc_config $name

: ${dracd_enable="NO"}

required_files="%%PREFIX%%/etc/dracd.allow"

run_rc_command "$1"
