# net/qadsl/pkg-plist
sbin/qadsl
etc/qadsl.conf.sample
etc/rc.d/qadsl.sh
