# sysutils/diskcheckd/files/openrc-diskcheckd.in
#!/sbin/openrc-run

name="diskcheckd"
command=%%PREFIX%%/sbin/diskcheckd

depend() {
	keyword -shutdown
}
