# www/zope210/files/zeo210.in
#!/bin/sh
#
# Startup script for Zeo server.
#
# $FreeBSD$
#

# PROVIDE: zeo210
# REQUIRE: DAEMON
# BEFORE: zope210

# Define these zeo210_* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#       /etc/rc.conf.d/zeo210
#
# zeo210_enable : bool
#   Enable Zeo ("YES") or not ("NO", the default).
#
# zeo210_instances : list
#   List of dirs with Zeo's instances ("" by default).
#

. /etc/rc.subr

name="zeo210"
rcvar=`set_rcvar`

zeo210ctl () {
    for instance in $zeo210_instances; do
	if [ -d ${instance} ]; then
	    echo -n "  Zeo instance ${instance} -> "
	    ${instance}/bin/zeoctl "$1"
	fi
    done
}

zeo210_start () {
    echo "Starting Zeo 2.10:"
    zeo210ctl "start"
}

zeo210_stop () {
    echo "Stopping Zeo 2.10:"
    zeo210ctl "stop"
}

zeo210_restart () {
    echo "Restarting Zeo 2.10:"
    zeo210ctl "restart"
}

  start_cmd="zeo210_start"
   stop_cmd="zeo210_stop"
restart_cmd="zeo210_restart"

load_rc_config $name

: ${zeo210_enable="NO"}
: ${zeo210_instances=""}

cmd="$1"
[ $# -gt 0 ] && shift
[ -n  "$*" ] && zeo210_instances="$*"

run_rc_command "${cmd}"
