# net/iplog/files/iplog.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: iplog
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# iplog_enable (bool):   Set to NO by default.
#               Set it to YES to enable iplog.
#

. /etc/rc.subr

name="iplog"
rcvar=${name}_enable

command=%%PREFIX%%/sbin/${name}
pidfile=/var/run/${name}/${name}.pid

load_rc_config $name

: ${iplog_enable="NO"}

command_args="-d -z"

run_rc_command "$1"
