# net-mgmt/zabbix22-server/files/openrc-zabbix_proxy.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf.local or /etc/rc.conf to
# enable zabbix_proxy:
#
# zabbix_proxy_config (string): Set to the standard config file path by
#         default.
# zabbix_proxy_pidfile (string):    Location of the zabbix_proxy pid file
#                               Default is /var/run/zabbix/zabbix_proxy.pid
# zabbix_proxy_paths (string): Set to standard path by default.  Set a search
#         if you have custom userparams that need binaries elsewhere.

name="zabbix_proxy"

: ${zabbix_proxy_config:=%%ETCDIR%%/${name}.conf}
: ${zabbix_proxy_pidfile:=/var/run/zabbix/zabbix_proxy.pid}

command="%%PREFIX%%/sbin/${name}"
required_files="${zabbix_proxy_config}"
pidfile=${zabbix_proxy_pidfile}

start_pre() {
	local rundir=${zabbix_proxy_pidfile%/*}
	if [ ! -d $rundir ] ; then
		install -d -m 0755 -o zabbix -g zabbix $rundir
	fi

	PATH=$zabbix_proxy_paths $command -c $zabbix_proxy_config
}

depend() {
	%%PGSQL%% need postgresql
	%%MYSQL%% need mysql-server
	keyword -shutdown
}
