# security/krb5-114/files/kpropd.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: kpropd
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# kpropd_enable (bool):      Set to NO by default.
#                            Set it to YES to enable kpropd.
# kpropd_flags (str):        Set to "" by default.

. /etc/rc.subr

name=kpropd
rcvar=kpropd_enable

load_rc_config $name

: ${kpropd_enable:="NO"}
: ${kpropd_flags=""}

command=%%PREFIX%%/sbin/${name}

run_rc_command "$1"
