# security/krb5-114/files/openrc-kpropd.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# kpropd_flags (str):        Set to "" by default.

name=kpropd

: ${kpropd_flags=""}

command=%%PREFIX%%/sbin/${name}

depend() {
	keyword -shutdown
}
