# sysutils/cbsd/pkg-plist
cbsd/zfs.subr
cbsd/upgrade/not_empty.txt
cbsd/tools/upgrade
cbsd/tools/unmountmd
cbsd/tools/unmountfstab
cbsd/tools/srcls
cbsd/tools/rsyncdon
cbsd/tools/rsyncdoff
cbsd/tools/repo-tui
cbsd/tools/repo
cbsd/tools/replacewdir
cbsd/tools/pwait.c
cbsd/tools/pkgbrowser
cbsd/tools/pkg_query
cbsd/tools/pkg_info
cbsd/tools/pkg_cutleaves
cbsd/tools/passwd
cbsd/tools/options-tui
cbsd/tools/objls
cbsd/tools/nodesql.h
cbsd/tools/nodesql.c
cbsd/tools/nics-list
cbsd/tools/naton
cbsd/tools/natoff
cbsd/tools/natcfg
cbsd/tools/mountmd
cbsd/tools/mountfstab
cbsd/tools/mkinventory
cbsd/tools/makescene
cbsd/tools/makejconf
cbsd/tools/makecsyncconf
cbsd/tools/kernells
cbsd/tools/jgensecrets
cbsd/tools/jconstruct-tui
cbsd/tools/jconstruct-fromargs
cbsd/tools/jconstruct
cbsd/tools/jailmapdb
cbsd/tools/jail2iso
cbsd/tools/ipfw-tables.c
cbsd/tools/inventorysql.h
cbsd/tools/inventorysql.c
cbsd/tools/initenv
cbsd/tools/imgpart
cbsd/tools/help
cbsd/tools/getnics-by-ip
cbsd/tools/getinfo
cbsd/tools/distccmakeconf
cbsd/tools/disks-list
cbsd/tools/dhcpd
cbsd/tools/detect-emulation
cbsd/tools/ddns
cbsd/tools/cpr-tui
cbsd/tools/cpr
cbsd/tools/copymini
cbsd/tools/checkrhost
cbsd/tools/checkip
cbsd/tools/cbsdsar
cbsd/tools/basels
cbsd/tools/adduser
cbsd/tools.subr
cbsd/system/world
cbsd/system/removesrc
cbsd/system/removeobj
cbsd/system/removekernel
cbsd/system/removebase
cbsd/system/preparebase
cbsd/system/kernel
cbsd/system/installworld
cbsd/system/installkernel
cbsd/system/delete-old-libs
cbsd/system/csup
cbsd/system/buildworld
cbsd/system/buildkernel
cbsd/system/base2pkg
cbsd/system.subr
cbsd/sudoexec/secretsfile
cbsd/sudoexec/ncshell
cbsd/sudoexec/mkdatadir
cbsd/sudoexec/jswmode
cbsd/sudoexec/jstop
cbsd/sudoexec/jstart
cbsd/sudoexec/jrsyncconf
cbsd/sudoexec/jlogin
cbsd/sudoexec/jaillogtail
cbsd/sudoexec/initenv
cbsd/sudoexec/check_ports
cbsd/sudoexec/cbsdchown
cbsd/share/rctl.conf
cbsd/share/minifiles9.txt
cbsd/share/minifiles10.txt
cbsd/share/jail-skel/resolv.conf
cbsd/share/jail-skel/rc.conf
cbsd/share/jail-skel/master.passwd
cbsd/share/jail-skel/make.conf
cbsd/share/jail-skel/crontab
cbsd/share/jail-arg
cbsd/share/defbase9.txt
cbsd/share/csync2.cfg
cbsd/share/cpr_packages.sh
cbsd/share/cpr_make.conf
cbsd/share/base2pkg/freebsd-zoneinfo-gcc-10.0.conf
cbsd/share/base2pkg/freebsd-zoneinfo-9.1.conf
cbsd/share/base2pkg/freebsd-zoneinfo-10.0.conf
cbsd/share/base2pkg/freebsd-src-gcc-10.0.conf
cbsd/share/base2pkg/freebsd-src-9.1.conf
cbsd/share/base2pkg/freebsd-src-10.0.conf
cbsd/share/base2pkg/freebsd-rescue-gcc-10.0.conf
cbsd/share/base2pkg/freebsd-rescue-9.1.conf
cbsd/share/base2pkg/freebsd-rescue-10.0.conf
cbsd/share/base2pkg/freebsd-ports-gcc-10.0.conf
cbsd/share/base2pkg/freebsd-ports-9.1.conf
cbsd/share/base2pkg/freebsd-ports-10.0.conf
cbsd/share/base2pkg/freebsd-lib32-gcc-10.0.conf
cbsd/share/base2pkg/freebsd-lib32-9.1.conf
cbsd/share/base2pkg/freebsd-lib32-10.0.conf
cbsd/share/base2pkg/freebsd-kernel-symbols-gcc-10.0.conf
cbsd/share/base2pkg/freebsd-kernel-symbols-9.1.conf
cbsd/share/base2pkg/freebsd-kernel-symbols-10.0.conf
cbsd/share/base2pkg/freebsd-kernel-gcc-10.0.conf
cbsd/share/base2pkg/freebsd-kernel-9.1.conf
cbsd/share/base2pkg/freebsd-kernel-10.0.conf
cbsd/share/base2pkg/freebsd-defaults-gcc-10.0.conf
cbsd/share/base2pkg/freebsd-defaults-9.1.conf
cbsd/share/base2pkg/freebsd-defaults-10.0.conf
cbsd/share/base2pkg/freebsd-base-gcc-10.0.conf
cbsd/share/base2pkg/freebsd-base-9.1.conf
cbsd/share/base2pkg/freebsd-base-10.0.conf
cbsd/settings-tui.subr
cbsd/securecmd
cbsd/sbin/nodeinetd.c
cbsd/sbin/netmask.c
cbsd/sbin/ncctld.h
cbsd/sbin/ncctld.c
cbsd/sbin/cbsdwdir.c
cbsd/sbin/cbsdtaskd/user.h
cbsd/sbin/cbsdtaskd/structs.h
cbsd/sbin/cbsdtaskd/pw_dup.h
cbsd/sbin/cbsdtaskd/popen.h
cbsd/sbin/cbsdtaskd/pathnames.h
cbsd/sbin/cbsdtaskd/misc.h
cbsd/sbin/cbsdtaskd/macros.h
cbsd/sbin/cbsdtaskd/job.h
cbsd/sbin/cbsdtaskd/globals.h
cbsd/sbin/cbsdtaskd/funcs.h
cbsd/sbin/cbsdtaskd/externs.h
cbsd/sbin/cbsdtaskd/env.h
cbsd/sbin/cbsdtaskd/entry.h
cbsd/sbin/cbsdtaskd/do_command.h
cbsd/sbin/cbsdtaskd/database.h
cbsd/sbin/cbsdtaskd/config.h
cbsd/sbin/cbsdtaskd/cbsdtaskd.h
cbsd/sbin/cbsdtaskd/cbsdtaskd.c
cbsd/sbin/cbsdtaskd/cbsdtask.h
cbsd/sbin/cbsdtaskd/cbsdtask.c
cbsd/sbin/cbsdtaskd/cbsdstuff.h
cbsd/sbin/cbsdtaskd/cbsdpath.h
cbsd/sbin/cbsdtaskd/bitstring.h
cbsd/sbin/cbsd.c
cbsd/release/mkisoimages.sh
cbsd/release/install.sh
cbsd/release/generate-release.sh
cbsd/rc.d/nodepinger
cbsd/rc.d/nodemsgd
cbsd/rc.d/jails-astart
cbsd/rc.d/cleartmp
cbsd/rc.d/cbsdtaskd
cbsd/rc.d/cbsdrsyncd
cbsd/rc.d/cbsd
cbsd/nodectl/rnode
cbsd/nodectl/nodetopinger
cbsd/nodectl/nodestatus
cbsd/nodectl/nodepingerls
cbsd/nodectl/nodeaddkey
cbsd/nodectl/node
cbsd/nodectl/modncping
cbsd/nc.subr
cbsd/misc/tmpfsdir
cbsd/misc/symlink_changer
cbsd/mdtools.subr
cbsd/man7/cbsd.7
cbsd/mailtools.subr
cbsd/jfs.subr
cbsd/jailctl/jwhereis
cbsd/jailctl/jupgrade
cbsd/jailctl/jswmode
cbsd/jailctl/jstop
cbsd/jailctl/jstatus
cbsd/jailctl/jstart
cbsd/jailctl/jsnapshot
cbsd/jailctl/jsetup-tui
cbsd/jailctl/jset
cbsd/jailctl/jrestart
cbsd/jailctl/jrename
cbsd/jailctl/jremove
cbsd/jailctl/jrctl-tui
cbsd/jailctl/jrctl
cbsd/jailctl/jrclone
cbsd/jailctl/jorder
cbsd/jailctl/jls
cbsd/jailctl/jlogin
cbsd/jailctl/jinfo
cbsd/jailctl/jimport
cbsd/jailctl/jget
cbsd/jailctl/jexport
cbsd/jailctl/jexec
cbsd/jailctl/jdescr
cbsd/jailctl/jcreate
cbsd/jailctl/jconfig
cbsd/jailctl/jcoldmigrate
cbsd/jailctl/jclone
cbsd/jailctl/jcleanup
cbsd/jailctl/jbackup
cbsd/jailctl/j2slave
cbsd/jailctl/j2prepare
cbsd/initenv.subr
cbsd/etc/zonekeys/not_empty.txt
cbsd/etc/src.conf
cbsd/etc/src-pxe.conf
cbsd/etc/src-lite.conf
cbsd/etc/src-gcc.conf
cbsd/etc/scenes/xrescuebsdkms/tmpfsdir
cbsd/etc/scenes/xrescuebsdkms/skel/root/bin/xconfig
cbsd/etc/scenes/xrescuebsdkms/skel/root/bin/conf/xorg.conf-patch2
cbsd/etc/scenes/xrescuebsdkms/skel/root/bin/conf/xorg.conf-patch
cbsd/etc/scenes/xrescuebsdkms/skel/root/.xinitrc
cbsd/etc/scenes/xrescuebsdkms/skel/root/.fluxbox/usermenu
cbsd/etc/scenes/xrescuebsdkms/skel/root/.fluxbox/overlay
cbsd/etc/scenes/xrescuebsdkms/skel/root/.fluxbox/menu
cbsd/etc/scenes/xrescuebsdkms/skel/root/.fluxbox/init
cbsd/etc/scenes/xrescuebsdkms/skel/root/.fluxbox/bsd_wallpaper-1280x960.jpg
cbsd/etc/scenes/xrescuebsdkms/skel/root/.cshrc
cbsd/etc/scenes/xrescuebsdkms/skel/root/.Xdefaults
cbsd/etc/scenes/xrescuebsdkms/skel/etc/ttys
cbsd/etc/scenes/xrescuebsdkms/skel/etc/spwd.db
cbsd/etc/scenes/xrescuebsdkms/skel/etc/rc.d/scons
cbsd/etc/scenes/xrescuebsdkms/skel/etc/rc.conf
cbsd/etc/scenes/xrescuebsdkms/skel/etc/pwd.db
cbsd/etc/scenes/xrescuebsdkms/skel/etc/passwd
cbsd/etc/scenes/xrescuebsdkms/skel/etc/motd
cbsd/etc/scenes/xrescuebsdkms/skel/etc/master.passwd
cbsd/etc/scenes/xrescuebsdkms/skel/etc/make.conf
cbsd/etc/scenes/xrescuebsdkms/skel/etc/group
cbsd/etc/scenes/xrescuebsdkms/skel/etc/csh.cshrc
cbsd/etc/scenes/xrescuebsdkms/descr
cbsd/etc/scenes/xrescuebsdkms/conf/xorg.conf-patch2
cbsd/etc/scenes/xrescuebsdkms/conf/xorg.conf-patch
cbsd/etc/scenes/xrescuebsdkms/bin/xconfig
cbsd/etc/scenes/xrescuebsdkms/bin/startx
cbsd/etc/scenes/xrescuebsdkms.scene
cbsd/etc/scenes/xrescuebsdkms.list
cbsd/etc/scenes/xrescuebsdkms.jconf
cbsd/etc/scenes/xrescuebsdkms-make.conf
cbsd/etc/scenes/wordpressngx/wordpressngx-descr
cbsd/etc/scenes/wordpressngx/skel/etc/resolv.conf
cbsd/etc/scenes/wordpressngx/skel/etc/rc.conf
cbsd/etc/scenes/wordpressngx/skel/etc/make.conf
cbsd/etc/scenes/wordpressngx.scene
cbsd/etc/scenes/wordpressngx.list
cbsd/etc/scenes/wordpressngx.jconf
cbsd/etc/scenes/squidproxy/descr
cbsd/etc/scenes/squidproxy.scene
cbsd/etc/scenes/squidproxy.list
cbsd/etc/scenes/squidproxy.jconf
cbsd/etc/scenes/rescuebsd/tmpfsdir
cbsd/etc/scenes/rescuebsd/skel/root/.cshrc
cbsd/etc/scenes/rescuebsd/skel/etc/ttys
cbsd/etc/scenes/rescuebsd/skel/etc/syslog.conf
cbsd/etc/scenes/rescuebsd/skel/etc/spwd.db
cbsd/etc/scenes/rescuebsd/skel/etc/rc.d/scons
cbsd/etc/scenes/rescuebsd/skel/etc/rc.conf
cbsd/etc/scenes/rescuebsd/skel/etc/pwd.db
cbsd/etc/scenes/rescuebsd/skel/etc/passwd
cbsd/etc/scenes/rescuebsd/skel/etc/motd
cbsd/etc/scenes/rescuebsd/skel/etc/master.passwd
cbsd/etc/scenes/rescuebsd/skel/etc/make.conf
cbsd/etc/scenes/rescuebsd/skel/etc/mail/aliases
cbsd/etc/scenes/rescuebsd/descr
cbsd/etc/scenes/rescuebsd.scene
cbsd/etc/scenes/rescuebsd.list
cbsd/etc/scenes/rescuebsd.jconf
cbsd/etc/scenes/rescuebsd-make.conf
cbsd/etc/scenes/redis/descr
cbsd/etc/scenes/redis.scene
cbsd/etc/scenes/redis.list
cbsd/etc/scenes/redis.jconf
cbsd/etc/scenes/postfixpgsql/descr
cbsd/etc/scenes/postfixpgsql.scene
cbsd/etc/scenes/postfixpgsql.list
cbsd/etc/scenes/postfixpgsql.jconf
cbsd/etc/scenes/nginx/descr
cbsd/etc/scenes/nginx.scene
cbsd/etc/scenes/nginx.list
cbsd/etc/scenes/nginx.jconf
cbsd/etc/scenes/nagioscacti/descr
cbsd/etc/scenes/nagioscacti.scene
cbsd/etc/scenes/nagioscacti.list
cbsd/etc/scenes/nagioscacti.jconf
cbsd/etc/scenes/mysql/descr
cbsd/etc/scenes/mysql.scene
cbsd/etc/scenes/mysql.list
cbsd/etc/scenes/mysql.jconf
cbsd/etc/scenes/mc/descr
cbsd/etc/scenes/mc.scene
cbsd/etc/scenes/mc.list
cbsd/etc/scenes/mc.jconf
cbsd/etc/scenes/joomlangx/descr
cbsd/etc/scenes/joomlangx.scene
cbsd/etc/scenes/joomlangx.list
cbsd/etc/scenes/joomlangx.jconf
cbsd/etc/scenes/git/descr
cbsd/etc/scenes/git.scene
cbsd/etc/scenes/git.list
cbsd/etc/scenes/git.jconf
cbsd/etc/scenes/drupal7ngx/descr
cbsd/etc/scenes/drupal7ngx.scene
cbsd/etc/scenes/drupal7ngx.list
cbsd/etc/scenes/drupal7ngx.jconf
cbsd/etc/scenes/baculasqlite/descr
cbsd/etc/scenes/baculasqlite.scene
cbsd/etc/scenes/baculasqlite.list
cbsd/etc/scenes/baculasqlite.jconf
cbsd/etc/scenes/asteriskpbx/descr
cbsd/etc/scenes/asteriskpbx.scene
cbsd/etc/scenes/asteriskpbx.list
cbsd/etc/scenes/asteriskpbx.jconf
cbsd/etc/scenes/amp/descr
cbsd/etc/scenes/amp.scene
cbsd/etc/scenes/amp.list
cbsd/etc/scenes/amp.jconf
cbsd/etc/pkg_cutleaves.conf
cbsd/etc/pkg.conf
cbsd/etc/periodic/not_empty.txt
cbsd/etc/makescene.conf
cbsd/etc/make-pxe.conf
cbsd/etc/make-dtrace.conf
cbsd/etc/lighttpd.conf
cbsd/etc/kernel-GENERIC-i386-9.1
cbsd/etc/kernel-GENERIC-i386-10.0
cbsd/etc/kernel-GENERIC-amd64-9.1
cbsd/etc/kernel-GENERIC-amd64-10.0
cbsd/etc/kernel-CBSD-amd64-10.0
cbsd/etc/devfs.rules
cbsd/etc/cpr.conf
cbsd/etc/cbsdnamed.conf
cbsd/etc/buildworld.conf
cbsd/etc/au_router/test2-router.cfg
cbsd/etc/au_router/test2-exclude.cfg
cbsd/etc/au_router/test1-router.cfg
cbsd/etc/au_router/test1-exclude.cfg
cbsd/etc/au_router/test1-client.cfg
cbsd/distccache.subr
cbsd/cbsd.conf
cbsd/build.subr
cbsd/bin/sqlite/sqlite3.h
cbsd/bin/sqlite/sqlite3.c
cbsd/bin/sqlite/shell.c
cbsd/bin/sqlite/Makefile
cbsd/bin/sbsd
cbsd/bin/rexe
cbsd/bin/nodepinger.c
cbsd/bin/cfetch.c
cbsd/bin/cbsdssh.c
cbsd/bin/cbsdsftp.c
cbsd/bin/cbsdlogtail.c
cbsd/bin/cbsd
cbsd/base2pkg.subr
cbsd/ansiicolor.subr
cbsd/Makefile
cbsd/.ssh/config
@dirrm cbsd/upgrade
@dirrm cbsd/tools
@dirrm cbsd/tmp
@dirrm cbsd/system
@dirrm cbsd/sudoexec
@dirrm cbsd/share/jail-skel
@dirrm cbsd/share/base2pkg
@dirrm cbsd/share
@dirrm cbsd/sbin/cbsdtaskd
@dirrm cbsd/sbin
@dirrm cbsd/release
@dirrm cbsd/rc.d
@dirrm cbsd/nodectl
@dirrm cbsd/misc
@dirrm cbsd/man7
@dirrm cbsd/jailctl
@dirrm cbsd/etc/zonekeys
@dirrm cbsd/etc/scenes/xrescuebsdkms/skel/root/bin/conf
@dirrm cbsd/etc/scenes/xrescuebsdkms/skel/root/bin
@dirrm cbsd/etc/scenes/xrescuebsdkms/skel/root/.fluxbox
@dirrm cbsd/etc/scenes/xrescuebsdkms/skel/root
@dirrm cbsd/etc/scenes/xrescuebsdkms/skel/etc/rc.d
@dirrm cbsd/etc/scenes/xrescuebsdkms/skel/etc
@dirrm cbsd/etc/scenes/xrescuebsdkms/skel
@dirrm cbsd/etc/scenes/xrescuebsdkms/conf
@dirrm cbsd/etc/scenes/xrescuebsdkms/bin
@dirrm cbsd/etc/scenes/xrescuebsdkms
@dirrm cbsd/etc/scenes/wordpressngx/skel/etc
@dirrm cbsd/etc/scenes/wordpressngx/skel
@dirrm cbsd/etc/scenes/wordpressngx
@dirrm cbsd/etc/scenes/squidproxy
@dirrm cbsd/etc/scenes/rescuebsd/skel/root
@dirrm cbsd/etc/scenes/rescuebsd/skel/etc/rc.d
@dirrm cbsd/etc/scenes/rescuebsd/skel/etc/mail
@dirrm cbsd/etc/scenes/rescuebsd/skel/etc
@dirrm cbsd/etc/scenes/rescuebsd/skel
@dirrm cbsd/etc/scenes/rescuebsd
@dirrm cbsd/etc/scenes/redis
@dirrm cbsd/etc/scenes/postfixpgsql
@dirrm cbsd/etc/scenes/nginx
@dirrm cbsd/etc/scenes/nagioscacti
@dirrm cbsd/etc/scenes/mysql
@dirrm cbsd/etc/scenes/mc
@dirrm cbsd/etc/scenes/joomlangx
@dirrm cbsd/etc/scenes/git
@dirrm cbsd/etc/scenes/drupal7ngx
@dirrm cbsd/etc/scenes/baculasqlite
@dirrm cbsd/etc/scenes/asteriskpbx
@dirrm cbsd/etc/scenes/amp
@dirrm cbsd/etc/scenes
@dirrm cbsd/etc/periodic
@dirrm cbsd/etc/au_router
@dirrm cbsd/etc
@dirrm cbsd/bin/sqlite
@dirrm cbsd/bin
@dirrm cbsd/.ssh
@dirrm cbsd
