# sysutils/anacron/files/anacron.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: anacron
# REQUIRE: LOGIN

#
# Add the following line to /etc/rc.conf to enable anacron:
#
# anacron_enable="YES"
#

. /etc/rc.subr

name=anacron
rcvar=anacron_enable

load_rc_config $name
: ${anacron_enable:=NO}
command="%%PREFIX%%/sbin/anacron"

run_rc_command "$1"
