# sysutils/anacron/files/anacron.in
#!/bin/sh

# PROVIDE: %%PORTNAME%%
# REQUIRE: LOGIN
# KEYWORD: resume

#
# Execute one of the following commands as root to enable %%PORTNAME%%:
#
#    service %%PORTNAME%% enable
#    sysrc %%PORTNAME%%_enable=YES
#

. /etc/rc.subr

name=%%PORTNAME%%
command="%%PREFIX%%/sbin/%%PORTNAME%%"
desc="%%COMMENT%%"
extra_commands=resume
rcvar=%%PORTNAME%%_enable
required_vars=syslogd_enable

load_rc_config $name
: ${%%PORTNAME%%_enable:=NO}

if [ $# -gt 0 -a $1 = "resume" ]; then
   arg=start
   if checkyesno $rcvar; then
       info "starting ${name} after resume"
   fi
else
   arg=$1
fi
debug "$name: \$1 = $1, arg = $arg"

run_rc_command "$arg"
