# x11-servers/Xfstt/pkg-plist
@comment $FreeBSD$
bin/xfstt
etc/rc.d/xfstt.sh
%%PORTDOCS%%%%DOCSDIR%%/CHANGES
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/THANKS
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@exec mkdir -p %D/lib/X11/fonts/TrueType
@unexec rm -f %D/lib/X11/fonts/TrueType/ttname.dir
@unexec rm -f %D/lib/X11/fonts/TrueType/ttinfo.dir
@unexec rmdir %D/lib/X11/fonts/TrueType 2>/dev/null || true
