# net/pppoa/pkg-plist
sbin/modem_run
sbin/pppoa2
sbin/pppoa3
etc/rc.d/adsl.sh.sample
libdata/mgmt.o
share/doc/speedtouch/FAQ
share/doc/speedtouch/FAQ-fr
share/doc/speedtouch/INSTALL
share/doc/speedtouch/INSTALL-de
share/doc/speedtouch/INSTALL-fr
share/doc/speedtouch/INSTALL-it
share/doc/speedtouch/INSTALL-nl
share/doc/speedtouch/ppp.conf.sample
share/doc/speedtouch/quick-install
@exec cat %D/share/doc/speedtouch/quick-install
share/speedtouch/boot.v123.bin
@dirrm share/doc/speedtouch
@dirrm share/speedtouch
