# emulators/vmware-guestd6/pkg-plist
@unexec %D/etc/rc.d/vmware-guestd.sh forcestop 2>/dev/null || true
etc/rc.d/vmware-guestd.sh
sbin/vmware-checkvm
sbin/vmware-guestd
%%VMWARE_VMMEMCTL:%%lib/vmware-tools/modules/vmmemctl.ko
%%VMWARE_VMXNET:%%lib/vmware-tools/modules/vmxnet.ko
%%VMWARE_KMODDIR:%%@dirrmtry lib/vmware-tools/modules
%%VMWARE_KMODDIR:%%@dirrmtry lib/vmware-tools
share/vmware-tools/poweroff-vm-default
share/vmware-tools/poweron-vm-default
share/vmware-tools/resume-vm-default
share/vmware-tools/suspend-vm-default
@dirrm share/vmware-tools
@cwd /
etc/vmware-tools
