# ftp/smbftpd/pkg-plist
@comment $FreeBSD$
etc/rc.d/smbftpd.sh
etc/smbftpd/smbftpd.conf
etc/smbftpd/smbftpd_share.conf
sbin/smbftpd
@dirrmtry etc/smbftpd
