# net/net-snmp4/pkg-plist
bin/encode_keychange
bin/snmpbulkget
bin/snmpbulkwalk
bin/snmpcheck
bin/snmpdelta
bin/snmpget
bin/snmpgetnext
bin/snmpinform
bin/snmpnetstat
bin/snmpset
bin/snmpstatus
bin/snmptable
bin/snmptest
bin/snmptranslate
bin/snmptrap
bin/snmpusm
bin/snmpwalk
bin/tkmib
etc/rc.d/snmpd.sh.sample
include/ucd-snmp/agent_read_config.h
include/ucd-snmp/agent_registry.h
include/ucd-snmp/agent_trap.h
include/ucd-snmp/asn1.h
include/ucd-snmp/auto_nlist.h
include/ucd-snmp/callback.h
include/ucd-snmp/default_store.h
include/ucd-snmp/ds_agent.h
include/ucd-snmp/header_complex.h
include/ucd-snmp/int64.h
include/ucd-snmp/keytools.h
include/ucd-snmp/mib.h
include/ucd-snmp/mibincl.h
include/ucd-snmp/parse.h
include/ucd-snmp/read_config.h
include/ucd-snmp/scapi.h
include/ucd-snmp/snmp.h
include/ucd-snmp/snmp_agent.h
include/ucd-snmp/snmp_alarm.h
include/ucd-snmp/snmp_api.h
include/ucd-snmp/snmp_client.h
include/ucd-snmp/snmp_debug.h
include/ucd-snmp/snmp_impl.h
include/ucd-snmp/snmp_logging.h
include/ucd-snmp/snmp_vars.h
include/ucd-snmp/snmpv3.h
include/ucd-snmp/struct.h
include/ucd-snmp/system.h
include/ucd-snmp/tools.h
include/ucd-snmp/transform_oids.h
include/ucd-snmp/ucd-snmp-agent-includes.h
include/ucd-snmp/ucd-snmp-config.h
include/ucd-snmp/ucd-snmp-includes.h
include/ucd-snmp/util_funcs.h
include/ucd-snmp/var_struct.h
include/ucd-snmp/version.h
@dirrm include/ucd-snmp
lib/libsnmp.so
lib/libsnmp.so.%%shlib%%
lib/libucdagent.so
lib/libucdagent.so.%%shlib%%
lib/libucdmibs.so
lib/libucdmibs.so.%%shlib%%
sbin/snmpd
sbin/snmptrapd
share/examples/ucd-snmp/EXAMPLE.conf
share/examples/ucd-snmp/ucdDemoPublic.conf
@dirrm share/examples/ucd-snmp
share/snmp/mibs/DISMAN-SCRIPT-MIB.txt
share/snmp/mibs/EtherLike-MIB.txt
share/snmp/mibs/HOST-RESOURCES-MIB.txt
share/snmp/mibs/HOST-RESOURCES-TYPES.txt
share/snmp/mibs/IANAifType-MIB.txt
share/snmp/mibs/IF-MIB.txt
share/snmp/mibs/IP-MIB.txt
share/snmp/mibs/IPV6-ICMP-MIB.txt
share/snmp/mibs/IPV6-MIB.txt
share/snmp/mibs/IPV6-TC.txt
share/snmp/mibs/IPV6-TCP-MIB.txt
share/snmp/mibs/IPV6-UDP-MIB.txt
share/snmp/mibs/RFC1155-SMI.txt
share/snmp/mibs/RFC1213-MIB.txt
share/snmp/mibs/RMON-MIB.txt
share/snmp/mibs/SNMP-COMMUNITY-MIB.txt
share/snmp/mibs/SNMP-FRAMEWORK-MIB.txt
share/snmp/mibs/SNMP-MPD-MIB.txt
share/snmp/mibs/SNMP-NOTIFICATION-MIB.txt
share/snmp/mibs/SNMP-PROXY-MIB.txt
share/snmp/mibs/SNMP-TARGET-MIB.txt
share/snmp/mibs/SNMP-USER-BASED-SM-MIB.txt
share/snmp/mibs/SNMP-VIEW-BASED-ACM-MIB.txt
share/snmp/mibs/SNMPv2-CONF.txt
share/snmp/mibs/SNMPv2-MIB.txt
share/snmp/mibs/SNMPv2-SMI.txt
share/snmp/mibs/SNMPv2-TC.txt
share/snmp/mibs/SNMPv2-TM.txt
share/snmp/mibs/TCP-MIB.txt
share/snmp/mibs/UCD-DEMO-MIB.txt
share/snmp/mibs/UCD-DISKIO-MIB.txt
share/snmp/mibs/UCD-DLMOD-MIB.txt
share/snmp/mibs/UCD-IPFWACC-MIB.txt
share/snmp/mibs/UCD-SNMP-MIB.txt
share/snmp/mibs/UDP-MIB.txt
@unexec rm -f %D/share/snmp/mibs/.index
@dirrm share/snmp/mibs
@dirrm share/snmp
@exec /usr/bin/env OBJFORMAT=%%PORTOBJFORMAT%% /sbin/ldconfig -m %D/lib
@unexec /usr/bin/env OBJFORMAT=%%PORTOBJFORMAT%% /sbin/ldconfig -R
