# net-im/ejabberd/pkg-plist
bin/ejabberd
bin/ejabberdctl
etc/ejabberd.cfg.sample
etc/ejabberd.defaults
etc/rc.d/ejabberd.sh
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ELDAPv3.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/acl.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/cyrsasl.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/cyrsasl_digest.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/cyrsasl_plain.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd.app
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_app.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_auth.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_auth_external.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_auth_internal.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_auth_ldap.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_auth_odbc.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_c2s.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_config.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_ctl.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_hooks.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_http.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_http_poll.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_listener.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_local.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_logger_h.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_receiver.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_router.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_s2s.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_s2s_in.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_s2s_out.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_service.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_sm.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_sup.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_tmp_sup.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_web.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/ejabberd_web_admin.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/eldap.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/extauth.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/gen_iq_handler.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/gen_mod.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/iconv.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/idna.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/jd2ejd.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/jlib.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_announce.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_configure.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_configure2.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_disco.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_echo.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_irc.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_irc_connection.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_last.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_last_odbc.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_muc.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_muc_room.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_offline.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_offline_odbc.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_privacy.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_private.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_pubsub.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_register.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_roster.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_roster_odbc.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_service_log.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_shared_roster.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_stats.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_time.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_vcard.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_vcard_ldap.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/mod_version.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/randoms.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/sha.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/shaper.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/stringprep.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/tls.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/translate.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/xml.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin/xml_stream.beam
lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/lib/expat_erl.so
lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/lib/iconv_erl.so
lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/lib/stringprep_drv.so
lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/lib/tls_drv.so
lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/msgs/de.msg
lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/msgs/es.msg
lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/msgs/fr.msg
lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/msgs/nl.msg
lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/msgs/pl.msg
lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/msgs/pt.msg
lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/msgs/ru.msg
lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/msgs/sv.msg
lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/msgs/uk.msg
%%PORTDOCS%%%%DOCSDIR%%/dev.aux
%%PORTDOCS%%%%DOCSDIR%%/dev.html
%%PORTDOCS%%%%DOCSDIR%%/dev.log
%%PORTDOCS%%%%DOCSDIR%%/dev.pdf
%%PORTDOCS%%%%DOCSDIR%%/dev.tex
%%PORTDOCS%%%%DOCSDIR%%/dev.toc
%%PORTDOCS%%%%DOCSDIR%%/disco.png
%%PORTDOCS%%%%DOCSDIR%%/discorus.png
%%PORTDOCS%%%%DOCSDIR%%/flow.dot
%%PORTDOCS%%%%DOCSDIR%%/guide.aux
%%PORTDOCS%%%%DOCSDIR%%/guide.html
%%PORTDOCS%%%%DOCSDIR%%/guide.log
%%PORTDOCS%%%%DOCSDIR%%/guide.out
%%PORTDOCS%%%%DOCSDIR%%/guide.pdf
%%PORTDOCS%%%%DOCSDIR%%/guide.tex
%%PORTDOCS%%%%DOCSDIR%%/guide.toc
%%PORTDOCS%%%%DOCSDIR%%/hevea.sty
%%PORTDOCS%%%%DOCSDIR%%/logo.png
%%PORTDOCS%%%%DOCSDIR%%/release_notes_%%PORTVERSION%%.txt
%%PORTDOCS%%%%DOCSDIR%%/release_notes_0.9.txt
%%PORTDOCS%%%%DOCSDIR%%/webadmmain.png
%%PORTDOCS%%%%DOCSDIR%%/webadmmainru.png
@dirrm %%PORTDOCS%%share/doc/ejabberd
@dirrm lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/msgs
@dirrm lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv/lib
@dirrm lib/erlang/lib/ejabberd-%%PORTVERSION%%/priv
@dirrm lib/erlang/lib/ejabberd-%%PORTVERSION%%/ebin
@dirrm lib/erlang/lib/ejabberd-%%PORTVERSION%%
