# net-im/ejabberd/pkg-plist
bin/ejabberd
bin/ejabberdctl
etc/ejabberd.cfg.sample
etc/rc.d/ejabberd.sh
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ELDAPv3.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/acl.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/configure.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/cyrsasl.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/cyrsasl_digest.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/cyrsasl_plain.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_app.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_auth.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_c2s.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_config.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_ctl.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_http.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_http_poll.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_listener.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_local.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_logger_h.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/eldap.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_receiver.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_router.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_s2s.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_s2s_in.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_s2s_out.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_service.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_sm.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_sup.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_tmp_sup.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_web.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd_web_admin.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/gen_iq_handler.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/gen_mod.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/iconv.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/idna.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/jd2ejd.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/jlib.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_configure.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_configure2.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_disco.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_echo.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_irc.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_irc_connection.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_last.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_muc.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_muc_room.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_offline.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_privacy.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_private.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_pubsub.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_register.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_roster.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_stats.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_time.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_vcard.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/mod_version.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/randoms.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/sha.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/shaper.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/stringprep.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/translate.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/xml.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/xml_stream.beam
lib/erlang/lib/ejabberd-%%VERSION%%/ebin/ejabberd.app
lib/erlang/lib/ejabberd-%%VERSION%%/priv/lib/expat_erl.so
lib/erlang/lib/ejabberd-%%VERSION%%/priv/lib/iconv_erl.so
lib/erlang/lib/ejabberd-%%VERSION%%/priv/lib/stringprep_drv.so
lib/erlang/lib/ejabberd-%%VERSION%%/priv/msgs/es.msg
lib/erlang/lib/ejabberd-%%VERSION%%/priv/msgs/fr.msg
lib/erlang/lib/ejabberd-%%VERSION%%/priv/msgs/nl.msg
lib/erlang/lib/ejabberd-%%VERSION%%/priv/msgs/ru.msg
%%PORTDOCS%%%%DOCSDIR%%/confacls.png
%%PORTDOCS%%%%DOCSDIR%%/confhostname.png
%%PORTDOCS%%%%DOCSDIR%%/dev.html
%%PORTDOCS%%%%DOCSDIR%%/dev.tex
%%PORTDOCS%%%%DOCSDIR%%/disco.png
%%PORTDOCS%%%%DOCSDIR%%/discoallusers.png
%%PORTDOCS%%%%DOCSDIR%%/discorunnodes.png
%%PORTDOCS%%%%DOCSDIR%%/discorus.png
%%PORTDOCS%%%%DOCSDIR%%/guide.html
%%PORTDOCS%%%%DOCSDIR%%/guide.tex
%%PORTDOCS%%%%DOCSDIR%%/logo.png
%%PORTDOCS%%%%DOCSDIR%%/webadmmain.png
%%PORTDOCS%%%%DOCSDIR%%/webadmmainru.png
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm lib/erlang/lib/ejabberd-%%VERSION%%/priv/lib
@dirrm lib/erlang/lib/ejabberd-%%VERSION%%/priv/msgs
@dirrm lib/erlang/lib/ejabberd-%%VERSION%%/priv
@dirrm lib/erlang/lib/ejabberd-%%VERSION%%/ebin
@dirrm lib/erlang/lib/ejabberd-%%VERSION%%
