# net/irrd/pkg-plist
sbin/irrd
sbin/irr_rpsl_submit
sbin/irr_check
sbin/irr_notify
sbin/irrdcacher
sbin/ripe2rpsl
sbin/update_source
etc/irrd.conf.sample
etc/rc.d/irrd.sh
