# sysutils/puppet4/files/openrc-puppetmaster.in
#!/sbin/openrc-run

name="puppetmaster"

: ${puppetmaster_rundir="/var/run/puppet"}
: ${puppetmaster_flags="--rundir=${puppetmaster_rundir}"}

command="%%PREFIX%%/bin/puppet"
command_args="master ${puppetmaster_flags}"
unset puppetmaster_flags

pidfile="${puppetmaster_rundir}/master.pid"
puppet_manifest="$($command config print manifest)"

start_pre() {
	checkconfig
}

restart_pre() {
	checkconfig
}

checkconfig() {
	einfo "Performing sanity check of ${name} configuration: "
	if ! ${command} parser validate "${puppet_manifest}"
	then
		einfo "FAILED"
		return 1
	else
		einfo "OK"
		return 0
	fi
}

if [ -n "${puppetmaster_mongrel_ports}" ]; then
	for port in ${puppetmaster_mongrel_ports}; do
		pidfile="${puppetmaster_rundir}/mongrel.${port}.pid"
		command_args="master --pidfile=${pidfile} --servertype=mongrel --masterport=${port}"
	done
fi
