# net/relayd/files/relayd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: relayd
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable relayd:
# relayd_enable="YES"
# relayd_flags="<set as needed>"

. /etc/rc.subr

name=relayd
rcvar=${name}_enable

load_rc_config $name

: ${relayd_enable="NO"}

command="%%PREFIX%%/sbin/relayd"
relayctl="%%PREFIX%%/sbin/relayctl"
reload_cmd="relayd_reload_cmd"
extra_commands="reload"

relayd_reload_cmd () {
	${relayctl} reload
}

run_rc_command "$1"
