# dns/dnsflood/files/dnsflood.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: dnsflood
# REQUIRE: NETWORKING
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable dnsflood:
#
# dnsflood_enable="YES"
#

. /etc/rc.subr

name=dnsflood
rcvar=`set_rcvar`

command=%%PREFIX%%/sbin/dns_flood_detector
command_args="-d"

# set defaults
dnsflood_enable=${dnsflood_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"
