# dns/dnsflood/files/openrc-dnsflood.in
#!/sbin/openrc-run

name=dnsflood

command=%%PREFIX%%/sbin/dns_flood_detector
command_args="-d"

depend() {
	use net
	keyword -shutdown
}
