# comms/telldus-core/files/telldusd.in
#!/bin/sh

#
# PROVIDE: telldusd
# Add the following lines to /etc/rc.conf to enable telldusd:
# telldusd_enable (bool):    Set to "NO" by default.
#                            Set it to "YES" to enable telldusd.
# telldusd_flags (str):      Set to "" by default.
# telldusd_configfile (str): Set to "/usr/local/etc/tellstick.conf" by defult
#

. /etc/rc.subr

name="telldusd"
rcvar=telldusd_enable

command="/usr/local/sbin/telldusd"
pidfile="/var/run/telldusd.pid"

telldusd_enable=${telldusd_enable:-"NO"}
telldusd_flags=${telldusd_flags:-""}
telldusd_configfile=${telldusd_configfile:-"/usr/local/etc/tellstick.conf"}

load_rc_config "${name}"

required_files="${telldusd_configfile}"

run_rc_command "$1"
