# audio/aureal-kmod/pkg-plist
etc/rc.d/aureal.sh.sample
@unexec /usr/bin/chflags noschg %D/lib/au88x0/snd_au8810.ko
lib/au88x0/snd_au8810.ko
@unexec /usr/bin/chflags noschg %D/lib/au88x0/snd_au8820.ko
lib/au88x0/snd_au8820.ko
@unexec /usr/bin/chflags noschg %D/lib/au88x0/snd_au8830.ko
lib/au88x0/snd_au8830.ko
@dirrm lib/au88x0
