# dns/register-dnsclient/pkg-plist
sbin/%%PORTNAME%%
@sample etc/%%PORTNAME%%.nsupdate.sample
etc/cron.d/%%PORTNAME%%
