# security/clamav-devel/pkg-plist
@comment $FreeBSD$
%%CLAMAV-MILTER%%@unexec [ ! -f %%RUNDIR%%/clamav-milter.pid ] || %%RC_DIR%%/clamav-milter%%RC_SUFX%% stop || true
@unexec [ ! -f %%RUNDIR%%/freshclam.pid ] || %%RC_DIR%%/clamav-freshclam%%RC_SUFX%% stop || true
@unexec [ ! -f %%RUNDIR%%/clamd.pid ] || %%RC_DIR%%/clamav-clamd%%RC_SUFX%% stop || true
bin/clamav-config
bin/clamscan
bin/clamdscan
bin/freshclam
bin/sigtool
etc/rc.d/clamav-clamd%%RC_SUFX%%
etc/rc.d/clamav-freshclam%%RC_SUFX%%
%%CLAMAV-MILTER%%etc/rc.d/clamav-milter%%RC_SUFX%%
%%CLAMAV-MILTER%%sbin/clamav-milter
sbin/clamd
include/clamav.h
lib/libclamav.so.1
lib/libclamav.so
lib/libclamav.a
libdata/pkgconfig/libclamav.pc
@unexec if cmp -s %D/etc/clamd.conf %D/etc/clamd.conf.default; then rm -f %D/etc/clamd.conf; fi
etc/clamd.conf.default
@exec [ -f %B/clamd.conf ] || cp %B/%f %B/clamd.conf
@unexec if cmp -s %D/etc/freshclam.conf %D/etc/freshclam.conf.default; then rm -f %D/etc/freshclam.conf; fi
etc/freshclam.conf.default
@exec [ -f %B/freshclam.conf ] || cp %B/%f %B/freshclam.conf
@unexec [ -s %%DBDIR%%/daily.cvd ] || rm -f %%DBDIR%%/daily.cvd || true
@unexec [ -s %%DBDIR%%/main.cvd ] || rm -f %%DBDIR%%/main.cvd || true
@unexec rmdir %%DBDIR%% 2>/dev/null || true
@unexec rmdir %%LOGDIR%% 2>/dev/null || true
@unexec rmdir %%RUNDIR%% 2>/dev/null || true
