# security/clamav-devel/pkg-plist
@exec mkdir -p share/clamav
bin/clamscan
bin/clamdscan
bin/freshclam
bin/sigtool
etc/rc.d/clamav-clamd.sh.sample
%%CLAMAV-MILTER:%%etc/rc.d/clamav-milter.sh.sample
%%CLAMAV-MILTER:%%sbin/clamav-milter
sbin/clamd
include/clamav.h
lib/libclamav.so.1
lib/libclamav.so
lib/libclamav.la
lib/libclamav.a
share/clamav/mirrors.txt
@dirrm share/clamav
@unexec echo ""
@unexec echo "========================================================="
@unexec echo " If you want to deinstall this package permanently then"
@unexec echo "               \"rmuser -y %%CLAMAVUSER%%\" "
@unexec echo " This will remove clamav's spool-directory and user "
@unexec echo "========================================================="
@unexec echo ""
