# security/clamav-devel/pkg-plist
bin/clamscan
bin/clamdscan
bin/freshclam
bin/sigtool
etc/rc.d/clamav-clamd%%RC_SUFX%%
%%CLAMAV-MILTER:%%etc/rc.d/clamav-milter%%RC_SUFX%%
%%CLAMAV-MILTER:%%sbin/clamav-milter
sbin/clamd
include/clamav.h
lib/libclamav.so.1
lib/libclamav.so
lib/libclamav.a
@unexec if cmp -s %D/etc/clamav.conf %D/etc/clamav.conf.default; then rm -f %D/etc/clamav.conf; fi
etc/clamav.conf.default
@exec [ -f %B/clamav.conf ] || cp %B/%f %B/clamav.conf
@owner %%CLAMAVUSER%%
@group %%CLAMAVGROUP%%
%%DATADIR%%/mirrors.txt
%%DATADIR%%/main.cvd
%%DATADIR%%/daily.cvd
@exec chown %%CLAMAVUSER%%:%%CLAMAVGROUP%% %B
@dirrm %%DATADIR%%
@unexec rmdir /var/run/clamav 2>/dev/null || true
@unexec rmdir /var/log/clamav 2>/dev/null || true
