# www/ot-recorder/pkg-plist
@sample(,ot-recorder,640) %%ETCDIR%%/ot-recorder.conf.sample
bin/ocat
etc/rc.d/ot-recorder
sbin/ot-recorder
%%WWWDIR%%/index.html
%%WWWDIR%%/last/index.html
%%WWWDIR%%/last/last.html
%%WWWDIR%%/last/manifest.json
%%WWWDIR%%/last/map_google.js
%%WWWDIR%%/last/map_leaflet.js
%%WWWDIR%%/last/websock.js
%%WWWDIR%%/map/index.html
%%WWWDIR%%/map/map_google.js
%%WWWDIR%%/map/map_leaflet.js
%%WWWDIR%%/owntracks.xsl
%%WWWDIR%%/static/datatables/css/jquery.dataTables.min.css
%%WWWDIR%%/static/datatables/images/sort_asc.png
%%WWWDIR%%/static/datatables/images/sort_asc_disabled.png
%%WWWDIR%%/static/datatables/images/sort_both.png
%%WWWDIR%%/static/datatables/images/sort_desc.png
%%WWWDIR%%/static/datatables/images/sort_desc_disabled.png
%%WWWDIR%%/static/datatables/js/jquery.dataTables.min.js
%%WWWDIR%%/static/datatables/js/jquery.min.js
%%WWWDIR%%/static/defaultface.svg
%%WWWDIR%%/static/index.html
%%WWWDIR%%/static/leaflet/images/layers-2x.png
%%WWWDIR%%/static/leaflet/images/layers.png
%%WWWDIR%%/static/leaflet/images/marker-icon-2x.png
%%WWWDIR%%/static/leaflet/images/marker-icon.png
%%WWWDIR%%/static/leaflet/images/marker-shadow.png
%%WWWDIR%%/static/leaflet/leaflet.css
%%WWWDIR%%/static/leaflet/leaflet.js
%%WWWDIR%%/static/leaflet/leaflet.js.map
%%WWWDIR%%/static/recorder.png
%%WWWDIR%%/static/flags/AD.png
%%WWWDIR%%/static/flags/AE.png
%%WWWDIR%%/static/flags/AF.png
%%WWWDIR%%/static/flags/AG.png
%%WWWDIR%%/static/flags/AI.png
%%WWWDIR%%/static/flags/AL.png
%%WWWDIR%%/static/flags/AM.png
%%WWWDIR%%/static/flags/AN.png
%%WWWDIR%%/static/flags/AO.png
%%WWWDIR%%/static/flags/AQ.png
%%WWWDIR%%/static/flags/AR.png
%%WWWDIR%%/static/flags/AS.png
%%WWWDIR%%/static/flags/AT.png
%%WWWDIR%%/static/flags/AU.png
%%WWWDIR%%/static/flags/AW.png
%%WWWDIR%%/static/flags/AX.png
%%WWWDIR%%/static/flags/AZ.png
%%WWWDIR%%/static/flags/BA.png
%%WWWDIR%%/static/flags/BB.png
%%WWWDIR%%/static/flags/BD.png
%%WWWDIR%%/static/flags/BE.png
%%WWWDIR%%/static/flags/BF.png
%%WWWDIR%%/static/flags/BG.png
%%WWWDIR%%/static/flags/BH.png
%%WWWDIR%%/static/flags/BI.png
%%WWWDIR%%/static/flags/BJ.png
%%WWWDIR%%/static/flags/BL.png
%%WWWDIR%%/static/flags/BM.png
%%WWWDIR%%/static/flags/BN.png
%%WWWDIR%%/static/flags/BO.png
%%WWWDIR%%/static/flags/BR.png
%%WWWDIR%%/static/flags/BS.png
%%WWWDIR%%/static/flags/BT.png
%%WWWDIR%%/static/flags/BW.png
%%WWWDIR%%/static/flags/BY.png
%%WWWDIR%%/static/flags/BZ.png
%%WWWDIR%%/static/flags/CA.png
%%WWWDIR%%/static/flags/CC.png
%%WWWDIR%%/static/flags/CD.png
%%WWWDIR%%/static/flags/CF.png
%%WWWDIR%%/static/flags/CG.png
%%WWWDIR%%/static/flags/CH.png
%%WWWDIR%%/static/flags/CI.png
%%WWWDIR%%/static/flags/CK.png
%%WWWDIR%%/static/flags/CL.png
%%WWWDIR%%/static/flags/CM.png
%%WWWDIR%%/static/flags/CN.png
%%WWWDIR%%/static/flags/CO.png
%%WWWDIR%%/static/flags/CR.png
%%WWWDIR%%/static/flags/CU.png
%%WWWDIR%%/static/flags/CV.png
%%WWWDIR%%/static/flags/CW.png
%%WWWDIR%%/static/flags/CX.png
%%WWWDIR%%/static/flags/CY.png
%%WWWDIR%%/static/flags/CZ.png
%%WWWDIR%%/static/flags/DE.png
%%WWWDIR%%/static/flags/DJ.png
%%WWWDIR%%/static/flags/DK.png
%%WWWDIR%%/static/flags/DM.png
%%WWWDIR%%/static/flags/DO.png
%%WWWDIR%%/static/flags/DZ.png
%%WWWDIR%%/static/flags/EC.png
%%WWWDIR%%/static/flags/EE.png
%%WWWDIR%%/static/flags/EG.png
%%WWWDIR%%/static/flags/EH.png
%%WWWDIR%%/static/flags/ER.png
%%WWWDIR%%/static/flags/ES.png
%%WWWDIR%%/static/flags/ET.png
%%WWWDIR%%/static/flags/EU.png
%%WWWDIR%%/static/flags/FI.png
%%WWWDIR%%/static/flags/FJ.png
%%WWWDIR%%/static/flags/FK.png
%%WWWDIR%%/static/flags/FM.png
%%WWWDIR%%/static/flags/FO.png
%%WWWDIR%%/static/flags/FR.png
%%WWWDIR%%/static/flags/GA.png
%%WWWDIR%%/static/flags/GB.png
%%WWWDIR%%/static/flags/GD.png
%%WWWDIR%%/static/flags/GE.png
%%WWWDIR%%/static/flags/GG.png
%%WWWDIR%%/static/flags/GH.png
%%WWWDIR%%/static/flags/GI.png
%%WWWDIR%%/static/flags/GL.png
%%WWWDIR%%/static/flags/GM.png
%%WWWDIR%%/static/flags/GN.png
%%WWWDIR%%/static/flags/GQ.png
%%WWWDIR%%/static/flags/GR.png
%%WWWDIR%%/static/flags/GS.png
%%WWWDIR%%/static/flags/GT.png
%%WWWDIR%%/static/flags/GU.png
%%WWWDIR%%/static/flags/GW.png
%%WWWDIR%%/static/flags/GY.png
%%WWWDIR%%/static/flags/HK.png
%%WWWDIR%%/static/flags/HN.png
%%WWWDIR%%/static/flags/HR.png
%%WWWDIR%%/static/flags/HT.png
%%WWWDIR%%/static/flags/HU.png
%%WWWDIR%%/static/flags/Hello.txt
%%WWWDIR%%/static/flags/IC.png
%%WWWDIR%%/static/flags/ID.png
%%WWWDIR%%/static/flags/IE.png
%%WWWDIR%%/static/flags/IL.png
%%WWWDIR%%/static/flags/IM.png
%%WWWDIR%%/static/flags/IN.png
%%WWWDIR%%/static/flags/IQ.png
%%WWWDIR%%/static/flags/IR.png
%%WWWDIR%%/static/flags/IS.png
%%WWWDIR%%/static/flags/IT.png
%%WWWDIR%%/static/flags/JE.png
%%WWWDIR%%/static/flags/JM.png
%%WWWDIR%%/static/flags/JO.png
%%WWWDIR%%/static/flags/JP.png
%%WWWDIR%%/static/flags/KE.png
%%WWWDIR%%/static/flags/KG.png
%%WWWDIR%%/static/flags/KH.png
%%WWWDIR%%/static/flags/KI.png
%%WWWDIR%%/static/flags/KM.png
%%WWWDIR%%/static/flags/KN.png
%%WWWDIR%%/static/flags/KP.png
%%WWWDIR%%/static/flags/KR.png
%%WWWDIR%%/static/flags/KW.png
%%WWWDIR%%/static/flags/KY.png
%%WWWDIR%%/static/flags/KZ.png
%%WWWDIR%%/static/flags/LA.png
%%WWWDIR%%/static/flags/LB.png
%%WWWDIR%%/static/flags/LC.png
%%WWWDIR%%/static/flags/LI.png
%%WWWDIR%%/static/flags/LICENSE.txt
%%WWWDIR%%/static/flags/LK.png
%%WWWDIR%%/static/flags/LR.png
%%WWWDIR%%/static/flags/LS.png
%%WWWDIR%%/static/flags/LT.png
%%WWWDIR%%/static/flags/LU.png
%%WWWDIR%%/static/flags/LV.png
%%WWWDIR%%/static/flags/LY.png
%%WWWDIR%%/static/flags/MA.png
%%WWWDIR%%/static/flags/MC.png
%%WWWDIR%%/static/flags/MD.png
%%WWWDIR%%/static/flags/ME.png
%%WWWDIR%%/static/flags/MF.png
%%WWWDIR%%/static/flags/MG.png
%%WWWDIR%%/static/flags/MH.png
%%WWWDIR%%/static/flags/MK.png
%%WWWDIR%%/static/flags/ML.png
%%WWWDIR%%/static/flags/MM.png
%%WWWDIR%%/static/flags/MN.png
%%WWWDIR%%/static/flags/MO.png
%%WWWDIR%%/static/flags/MP.png
%%WWWDIR%%/static/flags/MQ.png
%%WWWDIR%%/static/flags/MR.png
%%WWWDIR%%/static/flags/MS.png
%%WWWDIR%%/static/flags/MT.png
%%WWWDIR%%/static/flags/MU.png
%%WWWDIR%%/static/flags/MV.png
%%WWWDIR%%/static/flags/MW.png
%%WWWDIR%%/static/flags/MX.png
%%WWWDIR%%/static/flags/MY.png
%%WWWDIR%%/static/flags/MZ.png
%%WWWDIR%%/static/flags/NA.png
%%WWWDIR%%/static/flags/NC.png
%%WWWDIR%%/static/flags/NE.png
%%WWWDIR%%/static/flags/NF.png
%%WWWDIR%%/static/flags/NG.png
%%WWWDIR%%/static/flags/NI.png
%%WWWDIR%%/static/flags/NL.png
%%WWWDIR%%/static/flags/NO.png
%%WWWDIR%%/static/flags/NP.png
%%WWWDIR%%/static/flags/NR.png
%%WWWDIR%%/static/flags/NU.png
%%WWWDIR%%/static/flags/NZ.png
%%WWWDIR%%/static/flags/OM.png
%%WWWDIR%%/static/flags/PA.png
%%WWWDIR%%/static/flags/PE.png
%%WWWDIR%%/static/flags/PF.png
%%WWWDIR%%/static/flags/PG.png
%%WWWDIR%%/static/flags/PH.png
%%WWWDIR%%/static/flags/PK.png
%%WWWDIR%%/static/flags/PL.png
%%WWWDIR%%/static/flags/PN.png
%%WWWDIR%%/static/flags/PR.png
%%WWWDIR%%/static/flags/PS.png
%%WWWDIR%%/static/flags/PT.png
%%WWWDIR%%/static/flags/PW.png
%%WWWDIR%%/static/flags/PY.png
%%WWWDIR%%/static/flags/QA.png
%%WWWDIR%%/static/flags/RO.png
%%WWWDIR%%/static/flags/RS.png
%%WWWDIR%%/static/flags/RU.png
%%WWWDIR%%/static/flags/RW.png
%%WWWDIR%%/static/flags/SA.png
%%WWWDIR%%/static/flags/SB.png
%%WWWDIR%%/static/flags/SC.png
%%WWWDIR%%/static/flags/SD.png
%%WWWDIR%%/static/flags/SE.png
%%WWWDIR%%/static/flags/SG.png
%%WWWDIR%%/static/flags/SH.png
%%WWWDIR%%/static/flags/SI.png
%%WWWDIR%%/static/flags/SK.png
%%WWWDIR%%/static/flags/SL.png
%%WWWDIR%%/static/flags/SM.png
%%WWWDIR%%/static/flags/SN.png
%%WWWDIR%%/static/flags/SO.png
%%WWWDIR%%/static/flags/SR.png
%%WWWDIR%%/static/flags/SS.png
%%WWWDIR%%/static/flags/ST.png
%%WWWDIR%%/static/flags/SV.png
%%WWWDIR%%/static/flags/SY.png
%%WWWDIR%%/static/flags/SZ.png
%%WWWDIR%%/static/flags/TC.png
%%WWWDIR%%/static/flags/TD.png
%%WWWDIR%%/static/flags/TF.png
%%WWWDIR%%/static/flags/TG.png
%%WWWDIR%%/static/flags/TH.png
%%WWWDIR%%/static/flags/TJ.png
%%WWWDIR%%/static/flags/TK.png
%%WWWDIR%%/static/flags/TL.png
%%WWWDIR%%/static/flags/TM.png
%%WWWDIR%%/static/flags/TN.png
%%WWWDIR%%/static/flags/TO.png
%%WWWDIR%%/static/flags/TR.png
%%WWWDIR%%/static/flags/TT.png
%%WWWDIR%%/static/flags/TV.png
%%WWWDIR%%/static/flags/TW.png
%%WWWDIR%%/static/flags/TZ.png
%%WWWDIR%%/static/flags/UA.png
%%WWWDIR%%/static/flags/UG.png
%%WWWDIR%%/static/flags/US.png
%%WWWDIR%%/static/flags/UY.png
%%WWWDIR%%/static/flags/UZ.png
%%WWWDIR%%/static/flags/VA.png
%%WWWDIR%%/static/flags/VC.png
%%WWWDIR%%/static/flags/VE.png
%%WWWDIR%%/static/flags/VG.png
%%WWWDIR%%/static/flags/VI.png
%%WWWDIR%%/static/flags/VN.png
%%WWWDIR%%/static/flags/VU.png
%%WWWDIR%%/static/flags/WF.png
%%WWWDIR%%/static/flags/WS.png
%%WWWDIR%%/static/flags/YE.png
%%WWWDIR%%/static/flags/YT.png
%%WWWDIR%%/static/flags/ZA.png
%%WWWDIR%%/static/flags/ZM.png
%%WWWDIR%%/static/flags/ZW.png
%%WWWDIR%%/static/flags/__.png
%%WWWDIR%%/table/config.js
%%WWWDIR%%/table/index.html
%%WWWDIR%%/table/otable.css
%%WWWDIR%%/test-ws.html
%%WWWDIR%%/test.txt
%%WWWDIR%%/utils/debug.js
%%WWWDIR%%/utils/config.js
%%WWWDIR%%/utils/index.html
%%WWWDIR%%/utils/map.css
%%WWWDIR%%/utils/map.js
%%WWWDIR%%/utils/map_google.js
%%WWWDIR%%/utils/map_leaflet.js
%%WWWDIR%%/utils/misc.js
%%WWWDIR%%/utils/network.js
%%WWWDIR%%/views/jane.json.sample
%%WWWDIR%%/views/leafletmap.html
@dir %%WWWDIR%%/last
@dir %%WWWDIR%%/map
@dir %%WWWDIR%%/static/leaflet/images
@dir %%WWWDIR%%/static/leaflet
@dir %%WWWDIR%%/static
@dir %%WWWDIR%%/static/flags
@dir %%WWWDIR%%/table
@dir %%WWWDIR%%/views
@dir(ot-recorder,ot-recorder,) /var/log/ot-recorder
@dir(ot-recorder,ot-recorder,) /var/run/ot-recorder
@dir(ot-recorder,ot-recorder,750) /var/db/owntracks/recorder
@dir /var/db/owntracks
