# news/newscache/pkg-plist
bin/updatenews
etc/newscache.conf-dist
etc/newscache.auth-dist
etc/rc.d/newscache.sh
sbin/newscache
sbin/newscacheclean
