# net/socks5/pkg-plist
etc/rc.d/socks5.sh
bin/socks5
bin/runsocks
bin/stopsocks
bin/rarchie
bin/rfinger
bin/rftp
bin/rping
bin/rtelnet
bin/rtraceroute
bin/rwhois
include/socks.h
lib/libsocks5.a
lib/libsocks5_sh.so
%%PORTDOCS%%%%EXAMPLESDIR%%/README
%%PORTDOCS%%%%EXAMPLESDIR%%/multiple_server.example
%%PORTDOCS%%%%EXAMPLESDIR%%/server_2_multiple_server.example
%%PORTDOCS%%%%EXAMPLESDIR%%/single_homed_host.example
%%PORTDOCS%%%%EXAMPLESDIR%%/socks5.conf.dualhomed
%%PORTDOCS%%%%EXAMPLESDIR%%/socks5.conf.gssapi
%%PORTDOCS%%%%EXAMPLESDIR%%/socks5.conf.multipleservers
%%PORTDOCS%%%%EXAMPLESDIR%%/socks5.conf.server2server
%%PORTDOCS%%%%EXAMPLESDIR%%/socks5.conf.singlehomed
%%PORTDOCS%%@dirrm %%EXAMPLESDIR%%
%%PORTDOCS%%%%DOCSDIR%%/rfc1928.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc1929.txt
%%PORTDOCS%%%%DOCSDIR%%/rfc1961.txt
%%PORTDOCS%%%%DOCSDIR%%/socks.faq
%%PORTDOCS%%@dirrm %%DOCSDIR%%
