# misc/tempcontrol/pkg-plist
bin/tempcontrol
share/tempcontrol/temperaturecontrolrc.sample
etc/rc.d/temperaturecontrol.sh.sample
@dirrm share/tempcontrol
