# sysutils/ipa/pkg-plist
bin/ipa
bin/ipastat
etc/ipa.conf.default
etc/rc.d/ipa.sh.sample
