# sysutils/ipa/pkg-plist
bin/ipa
bin/ipactl
bin/ipastat
etc/rc.d/ipa.sh
include/ipa_mod.h
