# dns/nsd/pkg-plist
@comment $FreeBSD$
etc/rc.d/nsd.sh
etc/nsd/nsd.conf.sample
sbin/nsd
sbin/zonec
sbin/nsdc
sbin/nsd-notify
sbin/nsd-checkconf
sbin/nsd-patch
sbin/nsd-xfer
@dirrm etc/nsd
