# dns/nsd/pkg-plist
@comment $FreeBSD$
etc/rc.d/nsd.sh
etc/nsd/nsd.zones.sample
etc/nsd/nsdc.conf.sample
sbin/nsd
sbin/zonec
sbin/nsdc
sbin/nsd-notify
sbin/nsd-xfer
@dirrm etc/nsd
