# databases/postgresql91-server/pkg-plist
%%PG_PREFIX%%bin/postmaster
%%PG_PREFIX%%bin/postgres
%%PG_PREFIX%%bin/ecpg
%%PG_PREFIX%%bin/initdb
%%PG_PREFIX%%bin/initlocation
%%PG_PREFIX%%bin/ipcclean
%%PG_PREFIX%%bin/pg_ctl
%%PG_PREFIX%%bin/pg_dump
%%PG_PREFIX%%bin/pg_restore
%%PG_PREFIX%%bin/pg_encoding
%%PG_PREFIX%%bin/pg_dumpall
%%PG_PREFIX%%bin/pg_id
%%PG_PREFIX%%bin/pg_passwd
%%PG_PREFIX%%bin/psql
%%PG_PREFIX%%bin/createdb
%%PG_PREFIX%%bin/dropdb
%%PG_PREFIX%%bin/createuser
%%PG_PREFIX%%bin/dropuser
%%PG_PREFIX%%bin/droplang
%%PG_PREFIX%%bin/vacuumdb
%%PG_PREFIX%%bin/createlang
%%PG_PREFIX%%bin/pg_config
etc/rc.d/010.pgsql.sh
%%PG_INCLUDE%%/lib/dllist.h
@dirrm %%PG_INCLUDE%%/lib
%%PG_INCLUDE%%/libpq/pqcomm.h
%%PG_INCLUDE%%/libpq/libpq-fs.h
@dirrm %%PG_INCLUDE%%/libpq
%%PG_INCLUDE%%/libpq++/pgconnection.h
%%PG_INCLUDE%%/libpq++/pgdatabase.h
%%PG_INCLUDE%%/libpq++/pgtransdb.h
%%PG_INCLUDE%%/libpq++/pgcursordb.h
%%PG_INCLUDE%%/libpq++/pglobject.h
@dirrm %%PG_INCLUDE%%/libpq++
%%PG_INCLUDE%%/c.h
%%PG_INCLUDE%%/postgres_ext.h
%%PG_INCLUDE%%/postgres_fe.h
%%PG_INCLUDE%%/os.h
%%PG_INCLUDE%%/config.h
%%PG_INCLUDE%%/libpq-fe.h
%%PG_INCLUDE%%/libpq-int.h
%%PG_INCLUDE%%/pqexpbuffer.h
%%PG_INCLUDE%%/ecpgerrno.h
%%PG_INCLUDE%%/ecpglib.h
%%PG_INCLUDE%%/ecpgtype.h
%%PG_INCLUDE%%/sqlca.h
%%PG_INCLUDE%%/sql3types.h
%%PG_INCLUDE%%/libpgeasy.h
%%PG_INCLUDE%%/libpq++.h
@dirrm %%PG_INCLUDE%%
%%PG_PREFIX%%lib/libpq.so
%%PG_PREFIX%%lib/libpq.a
%%PG_PREFIX%%lib/libpq.so.2
%%PG_PREFIX%%lib/libecpg.so
%%PG_PREFIX%%lib/libecpg.a
%%PG_PREFIX%%lib/libecpg.so.3
%%PG_PREFIX%%lib/libpgeasy.so
%%PG_PREFIX%%lib/libpgeasy.a
%%PG_PREFIX%%lib/libpgeasy.so.2
%%PG_PREFIX%%lib/libpq++.so
%%PG_PREFIX%%lib/libpq++.a
%%PG_PREFIX%%lib/libpq++.so.3
%%PG_PREFIX%%lib/plpgsql.so
share/postgresql/global.bki
share/postgresql/global.description
share/postgresql/post-install-notes
share/postgresql/template1.bki
share/postgresql/template1.description
share/postgresql/pg_hba.conf.sample
share/postgresql/pg_ident.conf.sample
share/postgresql/postgresql.conf.sample
@unexec if cmp -s ~pgsql/.profile %D/share/postgresql/dot.profile.dist; then rm -f ~pgsql/.profile; fi
@unexec if cmp -s ~pgsql/.cshrc %D/share/postgresql/dot.cshrc.dist; then rm -f ~pgsql/.cshrc; fi
@unexec [ -f ~pgsql/dot.profile.dist ] && rm ~pgsql/dot.profile.dist
@unexec [ -f ~pgsql/dot.cshrc.dist ] && rm ~pgsql/dot.cshrc.dist
@exec [ -f %D/share/postgresql/post-install-notes ] && cp %D/share/postgresql/post-install-notes ~pgsql && chown pgsql:pgsql ~pgsql/post-install-notes
@unexec [ -f ~pgsql/post-install-notes ] && rm ~pgsql/post-install-notes
share/postgresql/dot.profile.dist
share/postgresql/dot.cshrc.dist
@dirrm share/postgresql
@exec [ -f %D/share/postgresql/dot.profile.dist ] && cp %D/share/postgresql/dot.profile.dist ~pgsql/.
@exec [ -f %D/share/postgresql/dot.cshrc.dist ] && cp %D/share/postgresql/dot.cshrc.dist ~pgsql/.
@exec [ -f ~pgsql/.profile ] || (cp %D/share/postgresql/dot.profile.dist ~pgsql/.profile && chown pgsql:pgsql ~pgsql/.profile)
@exec [ -f ~pgsql/.cshrc ] || (cp %D/share/postgresql/dot.cshrc.dist ~pgsql/.cshrc && chown pgsql:pgsql ~pgsql/.cshrc)
@unexec rmdir %D/%%PG_PREFIX%%bin 2> /dev/null || true
@unexec rmdir %D/%%PG_PREFIX%%lib 2> /dev/null || true
@unexec rmdir %D/%%PG_PREFIX%%man/man1 2> /dev/null || true
@unexec rmdir %D/%%PG_PREFIX%%man/manl 2> /dev/null || true
@unexec rmdir %D/%%PG_PREFIX%%man 2> /dev/null || true
@unexec echo "Warning: If you will *NOT* use this package anymore, please remove %D/pgsql and the pgsql user manually."
