# security/pf/pkg-plist
libexec/ftp-proxy

modules/linker.hints
modules/pf.ko
modules/pflog.ko
modules/pfsync.ko
%%WITH_ALTQ%%modules/pfaltq.ko

sbin/pfctl
sbin/pflogd
sbin/pftcpdump

etc/rc.d/pf.sh.sample
etc/pf.conf.default
etc/pf.os

include/pf/net/pfvar.h
include/pf/net/if_pflog.h
include/pf/net/if_pfsync.h
%%WITH_ALTQ%%include/pf/altq/altq.h
%%WITH_ALTQ%%include/pf/altq/altq_cbq.h
%%WITH_ALTQ%%include/pf/altq/altq_cdnr.h
%%WITH_ALTQ%%include/pf/altq/altq_classq.h
%%WITH_ALTQ%%include/pf/altq/altq_hfsc.h
%%WITH_ALTQ%%include/pf/altq/altq_priq.h
%%WITH_ALTQ%%include/pf/altq/altq_red.h
%%WITH_ALTQ%%include/pf/altq/altq_rio.h
%%WITH_ALTQ%%include/pf/altq/altq_rmclass.h
%%WITH_ALTQ%%include/pf/altq/altq_rmclass_debug.h
%%WITH_ALTQ%%include/pf/altq/altq_var.h
%%WITH_ALTQ%%include/pf/altq/if_altq.h

@dirrm modules
@dirrm include/pf/net
@dirrm include/pf
