# audio/darkice/files/darkice.in
#!/bin/sh

# PROVIDE: darkice
# REQUIRE: DAEMON
# BEFORE:  LOGIN
# KEYWORD: SHUTDOWN
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# darkice_enable (bool):	Set to NO by default.
#				Set it to YES to enable darkice.

. /etc/rc.subr

name="darkice"
rcvar=darkice_enable

command="%%PREFIX%%/bin/$name"
command_args="& 1>/dev/null"
required_files="%%ETCDIR%%.cfg"

# read configuration and set defaults
load_rc_config $name

: ${darkice_enable="NO"}
: ${darkice_flags="-c ${required_files}"}

run_rc_command $1
