# mail/p5-Mail-SpamAssassin/pkg-plist
bin/sa-learn
bin/spamassassin
bin/spamc
bin/spamd
etc/mail/spamassassin/local.cf.sample
etc/rc.d/spamd.sh-dist
%%SITE_PERL%%/Mail/SpamAssassin.pm
%%SITE_PERL%%/Mail/SpamAssassin/ArchiveIterator.pm
%%SITE_PERL%%/Mail/SpamAssassin/AuditMessage.pm
%%SITE_PERL%%/Mail/SpamAssassin/AutoWhitelist.pm
%%SITE_PERL%%/Mail/SpamAssassin/Bayes.pm
%%SITE_PERL%%/Mail/SpamAssassin/BayesStore.pm
%%SITE_PERL%%/Mail/SpamAssassin/CmdLearn.pm
%%SITE_PERL%%/Mail/SpamAssassin/Conf.pm
%%SITE_PERL%%/Mail/SpamAssassin/ConfSourceSQL.pm
%%SITE_PERL%%/Mail/SpamAssassin/DBBasedAddrList.pm
%%SITE_PERL%%/Mail/SpamAssassin/Dns.pm
%%SITE_PERL%%/Mail/SpamAssassin/EncappedMIME.pm
%%SITE_PERL%%/Mail/SpamAssassin/EncappedMessage.pm
%%SITE_PERL%%/Mail/SpamAssassin/EvalTests.pm
%%SITE_PERL%%/Mail/SpamAssassin/HTML.pm
%%SITE_PERL%%/Mail/SpamAssassin/Locales.pm
%%SITE_PERL%%/Mail/SpamAssassin/Locker.pm
%%SITE_PERL%%/Mail/SpamAssassin/MailingList.pm
%%SITE_PERL%%/Mail/SpamAssassin/Message.pm
%%SITE_PERL%%/Mail/SpamAssassin/NetSet.pm
%%SITE_PERL%%/Mail/SpamAssassin/NoMailAudit.pm
%%SITE_PERL%%/Mail/SpamAssassin/PerMsgLearner.pm
%%SITE_PERL%%/Mail/SpamAssassin/PerMsgStatus.pm
%%SITE_PERL%%/Mail/SpamAssassin/PersistentAddrList.pm
%%SITE_PERL%%/Mail/SpamAssassin/PhraseFreqs.pm
%%SITE_PERL%%/Mail/SpamAssassin/Received.pm
%%SITE_PERL%%/Mail/SpamAssassin/Replier.pm
%%SITE_PERL%%/Mail/SpamAssassin/Reporter.pm
%%SITE_PERL%%/Mail/SpamAssassin/SHA1.pm
%%SITE_PERL%%/Mail/SpamAssassin/TextCat.pm
%%SITE_PERL%%/Mail/SpamAssassin/UnixLocker.pm
%%SITE_PERL%%/Mail/SpamAssassin/Util.pm
%%SITE_PERL%%/Mail/SpamAssassin/Win32Locker.pm
%%SITE_PERL%%/%%PERL_ARCH%%/auto/Mail/SpamAssassin/.packlist
%%PORTDOCS%%%%DOCSDIR%%/BUGS
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT
%%PORTDOCS%%%%DOCSDIR%%/Changes
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README.spamd
%%PORTDOCS%%%%DOCSDIR%%/README.spamd-vpopmail
%%PORTDOCS%%%%DOCSDIR%%/README.sql
%%PORTDOCS%%%%DOCSDIR%%/TRADEMARK
%%PORTDOCS%%%%DOCSDIR%%/USAGE
%%PORTDOCS%%%%DOCSDIR%%/procmailrc.example
%%PORTDOCS%%%%DOCSDIR%%/sample-nonspam.txt
%%PORTDOCS%%%%DOCSDIR%%/sample-spam.txt
%%PORTDOCS%%%%DOCSDIR%%/spamassassin.sql
%%PORTDOCS%%@dirrm %%DOCSDIR%%
share/spamassassin/10_misc.cf
share/spamassassin/20_anti_ratware.cf
share/spamassassin/20_body_tests.cf
share/spamassassin/20_compensate.cf
share/spamassassin/20_dnsbl_tests.cf
share/spamassassin/20_fake_helo_tests.cf
share/spamassassin/20_head_tests.cf
share/spamassassin/20_html_tests.cf
share/spamassassin/20_meta_tests.cf
share/spamassassin/20_phrases.cf
share/spamassassin/20_porn.cf
share/spamassassin/20_ratware.cf
share/spamassassin/20_uri_tests.cf
share/spamassassin/23_bayes.cf
share/spamassassin/25_body_tests_es.cf
share/spamassassin/25_body_tests_pl.cf
share/spamassassin/25_head_tests_es.cf
share/spamassassin/25_head_tests_pl.cf
share/spamassassin/30_text_de.cf
share/spamassassin/30_text_es.cf
share/spamassassin/30_text_fr.cf
share/spamassassin/30_text_it.cf
share/spamassassin/30_text_pl.cf
share/spamassassin/30_text_sk.cf
share/spamassassin/50_scores.cf
share/spamassassin/60_whitelist.cf
share/spamassassin/languages
share/spamassassin/triplets.txt
share/spamassassin/user_prefs.template
@dirrm %%SITE_PERL%%/%%PERL_ARCH%%/auto/Mail/SpamAssassin
@unexec rmdir %D/%%SITE_PERL%%/%%PERL_ARCH%%/auto/Mail 2>/dev/null || true
@dirrm %%SITE_PERL%%/Mail/SpamAssassin
@unexec rmdir %D/%%SITE_PERL%%/Mail 2>/dev/null || true
@unexec rmdir %D/etc/mail/spamassassin 2>/dev/null || true
@unexec rmdir %D/etc/mail 2>/dev/null || true
@dirrm share/spamassassin
