# mail/p5-Mail-SpamAssassin/pkg-plist
@unexec %D/etc/rc.d/sa-spamd.sh stop > /dev/null 2>&1 || true
bin/sa-learn
bin/spamassassin
bin/spamc
bin/spamd
etc/mail/spamassassin/local.cf.sample
@unexec if cmp -s %B/init.pre.sample %B/init.pre; then rm -f %B/init.pre; fi
etc/mail/spamassassin/init.pre.sample
@exec [ -f %B/init.pre ] || cp %B/%f %B/init.pre
etc/rc.d/sa-spamd.sh
%%SITE_PERL%%/Mail/SpamAssassin.pm
%%SITE_PERL%%/Mail/SpamAssassin/ArchiveIterator.pm
%%SITE_PERL%%/Mail/SpamAssassin/AutoWhitelist.pm
%%SITE_PERL%%/Mail/SpamAssassin/Bayes.pm
%%SITE_PERL%%/Mail/SpamAssassin/BayesStore.pm
%%SITE_PERL%%/Mail/SpamAssassin/BayesStore/DBM.pm
%%SITE_PERL%%/Mail/SpamAssassin/BayesStore/SQL.pm
%%SITE_PERL%%/Mail/SpamAssassin/Conf.pm
%%SITE_PERL%%/Mail/SpamAssassin/Conf/LDAP.pm
%%SITE_PERL%%/Mail/SpamAssassin/Conf/Parser.pm
%%SITE_PERL%%/Mail/SpamAssassin/Conf/SQL.pm
%%SITE_PERL%%/Mail/SpamAssassin/Constants.pm
%%SITE_PERL%%/Mail/SpamAssassin/DBBasedAddrList.pm
%%SITE_PERL%%/Mail/SpamAssassin/Dns.pm
%%SITE_PERL%%/Mail/SpamAssassin/EvalTests.pm
%%SITE_PERL%%/Mail/SpamAssassin/HTML.pm
%%SITE_PERL%%/Mail/SpamAssassin/Locales.pm
%%SITE_PERL%%/Mail/SpamAssassin/Locker.pm
%%SITE_PERL%%/Mail/SpamAssassin/Locker/Flock.pm
%%SITE_PERL%%/Mail/SpamAssassin/Locker/UnixNFSSafe.pm
%%SITE_PERL%%/Mail/SpamAssassin/Locker/Win32.pm
%%SITE_PERL%%/Mail/SpamAssassin/MailingList.pm
%%SITE_PERL%%/Mail/SpamAssassin/Message.pm
%%SITE_PERL%%/Mail/SpamAssassin/Message/Metadata.pm
%%SITE_PERL%%/Mail/SpamAssassin/Message/Metadata/Received.pm
%%SITE_PERL%%/Mail/SpamAssassin/Message/Node.pm
%%SITE_PERL%%/Mail/SpamAssassin/NetSet.pm
%%SITE_PERL%%/Mail/SpamAssassin/PerMsgLearner.pm
%%SITE_PERL%%/Mail/SpamAssassin/PerMsgStatus.pm
%%SITE_PERL%%/Mail/SpamAssassin/PersistentAddrList.pm
%%SITE_PERL%%/Mail/SpamAssassin/Plugin.pm
%%SITE_PERL%%/Mail/SpamAssassin/Plugin/Hashcash.pm
%%SITE_PERL%%/Mail/SpamAssassin/Plugin/RelayCountry.pm
%%SITE_PERL%%/Mail/SpamAssassin/Plugin/SPF.pm
%%SITE_PERL%%/Mail/SpamAssassin/Plugin/Test.pm
%%SITE_PERL%%/Mail/SpamAssassin/Plugin/URIDNSBL.pm
%%SITE_PERL%%/Mail/SpamAssassin/PluginHandler.pm
%%SITE_PERL%%/Mail/SpamAssassin/Reporter.pm
%%SITE_PERL%%/Mail/SpamAssassin/SQLBasedAddrList.pm
%%SITE_PERL%%/Mail/SpamAssassin/TextCat.pm
%%SITE_PERL%%/Mail/SpamAssassin/Util.pm
%%SITE_PERL%%/Mail/SpamAssassin/Util/RegistrarBoundaries.pm
%%SITE_PERL%%/%%PERL_ARCH%%/auto/Mail/SpamAssassin/.packlist
%%DATADIR%%/10_misc.cf
%%DATADIR%%/20_anti_ratware.cf
%%DATADIR%%/20_body_tests.cf
%%DATADIR%%/20_compensate.cf
%%DATADIR%%/20_dnsbl_tests.cf
%%DATADIR%%/20_drugs.cf
%%DATADIR%%/20_fake_helo_tests.cf
%%DATADIR%%/20_head_tests.cf
%%DATADIR%%/20_html_tests.cf
%%DATADIR%%/20_meta_tests.cf
%%DATADIR%%/20_phrases.cf
%%DATADIR%%/20_porn.cf
%%DATADIR%%/20_ratware.cf
%%DATADIR%%/20_uri_tests.cf
%%DATADIR%%/23_bayes.cf
%%DATADIR%%/25_body_tests_es.cf
%%DATADIR%%/25_hashcash.cf
%%DATADIR%%/25_spf.cf
%%DATADIR%%/25_uribl.cf
%%DATADIR%%/30_text_de.cf
%%DATADIR%%/30_text_fr.cf
%%DATADIR%%/30_text_nl.cf
%%DATADIR%%/30_text_pl.cf
%%DATADIR%%/50_scores.cf
%%DATADIR%%/60_whitelist.cf
%%DATADIR%%/languages
%%DATADIR%%/triplets.txt
%%DATADIR%%/user_prefs.template
%%TOOLS%%%%DATADIR%%/tools/README.speedtest
%%TOOLS%%%%DATADIR%%/tools/bayes_dump_to_trusted_networks
%%TOOLS%%%%DATADIR%%/tools/check_whitelist
%%TOOLS%%%%DATADIR%%/tools/convert_awl_dbm_to_sql
%%TOOLS%%%%DATADIR%%/tools/desc_length.pl
%%TOOLS%%%%DATADIR%%/tools/mboxsplit
%%TOOLS%%%%DATADIR%%/tools/sa-stats.pl
%%TOOLS%%%%DATADIR%%/tools/speedtest
%%TOOLS%%%%DATADIR%%/tools/split_corpora
%%TOOLS%%%%DATADIR%%/tools/test_extract
%%TOOLS%%%%DATADIR%%/tools/triplets.pl
%%TOOLS%%@dirrm %%DATADIR%%/tools
@dirrm %%DATADIR%%
@dirrm %%SITE_PERL%%/%%PERL_ARCH%%/auto/Mail/SpamAssassin
@unexec rmdir %D/%%SITE_PERL%%/%%PERL_ARCH%%/auto/Mail 2>/dev/null || true
@dirrm %%SITE_PERL%%/Mail/SpamAssassin/Util
@dirrm %%SITE_PERL%%/Mail/SpamAssassin/Plugin
@dirrm %%SITE_PERL%%/Mail/SpamAssassin/Message/Metadata
@dirrm %%SITE_PERL%%/Mail/SpamAssassin/Message
@dirrm %%SITE_PERL%%/Mail/SpamAssassin/Locker
@dirrm %%SITE_PERL%%/Mail/SpamAssassin/Conf
@dirrm %%SITE_PERL%%/Mail/SpamAssassin/BayesStore
@dirrm %%SITE_PERL%%/Mail/SpamAssassin
@unexec rmdir %D/%%SITE_PERL%%/Mail 2>/dev/null || true
@unexec rmdir %D/etc/mail/spamassassin 2>/dev/null || true
@unexec rmdir %D/etc/mail 2>/dev/null || true
@unexec rmdir /var/run/spamd 2>/dev/null || true
