# mail/p5-Mail-SpamAssassin/pkg-plist
bin/sa-learn
bin/spamassassin
bin/spamc
bin/spamd
etc/mail/spamassassin/local.cf.sample
etc/rc.d/spamd.sh
%%SITE_PERL%%/Mail/SpamAssassin.pm
%%SITE_PERL%%/Mail/SpamAssassin/ArchiveIterator.pm
%%SITE_PERL%%/Mail/SpamAssassin/AuditMessage.pm
%%SITE_PERL%%/Mail/SpamAssassin/AutoWhitelist.pm
%%SITE_PERL%%/Mail/SpamAssassin/Bayes.pm
%%SITE_PERL%%/Mail/SpamAssassin/BayesStore.pm
%%SITE_PERL%%/Mail/SpamAssassin/CmdLearn.pm
%%SITE_PERL%%/Mail/SpamAssassin/Conf.pm
%%SITE_PERL%%/Mail/SpamAssassin/ConfSourceSQL.pm
%%SITE_PERL%%/Mail/SpamAssassin/DBBasedAddrList.pm
%%SITE_PERL%%/Mail/SpamAssassin/Dns.pm
%%SITE_PERL%%/Mail/SpamAssassin/EncappedMIME.pm
%%SITE_PERL%%/Mail/SpamAssassin/EncappedMessage.pm
%%SITE_PERL%%/Mail/SpamAssassin/EvalTests.pm
%%SITE_PERL%%/Mail/SpamAssassin/HTML.pm
%%SITE_PERL%%/Mail/SpamAssassin/Locales.pm
%%SITE_PERL%%/Mail/SpamAssassin/Locker.pm
%%SITE_PERL%%/Mail/SpamAssassin/MailingList.pm
%%SITE_PERL%%/Mail/SpamAssassin/Message.pm
%%SITE_PERL%%/Mail/SpamAssassin/NetSet.pm
%%SITE_PERL%%/Mail/SpamAssassin/NoMailAudit.pm
%%SITE_PERL%%/Mail/SpamAssassin/PerMsgLearner.pm
%%SITE_PERL%%/Mail/SpamAssassin/PerMsgStatus.pm
%%SITE_PERL%%/Mail/SpamAssassin/PersistentAddrList.pm
%%SITE_PERL%%/Mail/SpamAssassin/PhraseFreqs.pm
%%SITE_PERL%%/Mail/SpamAssassin/Received.pm
%%SITE_PERL%%/Mail/SpamAssassin/Replier.pm
%%SITE_PERL%%/Mail/SpamAssassin/Reporter.pm
%%SITE_PERL%%/Mail/SpamAssassin/SHA1.pm
%%SITE_PERL%%/Mail/SpamAssassin/TextCat.pm
%%SITE_PERL%%/Mail/SpamAssassin/UnixLocker.pm
%%SITE_PERL%%/Mail/SpamAssassin/Util.pm
%%SITE_PERL%%/Mail/SpamAssassin/Win32Locker.pm
%%SITE_PERL%%/%%PERL_ARCH%%/auto/Mail/SpamAssassin/.packlist
%%PORTDOCS%%%%DOCSDIR%%/BUGS
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT
%%PORTDOCS%%%%DOCSDIR%%/Changes
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README.spamd
%%PORTDOCS%%%%DOCSDIR%%/README.spamd-vpopmail
%%PORTDOCS%%%%DOCSDIR%%/README.sql
%%PORTDOCS%%%%DOCSDIR%%/TRADEMARK
%%PORTDOCS%%%%DOCSDIR%%/USAGE
%%PORTDOCS%%%%DOCSDIR%%/procmailrc.example
%%PORTDOCS%%%%DOCSDIR%%/sample-nonspam.txt
%%PORTDOCS%%%%DOCSDIR%%/sample-spam.txt
%%PORTDOCS%%%%DOCSDIR%%/spamassassin.sql
%%DATADIR%%/10_misc.cf
%%DATADIR%%/20_anti_ratware.cf
%%DATADIR%%/20_body_tests.cf
%%DATADIR%%/20_compensate.cf
%%DATADIR%%/20_dnsbl_tests.cf
%%DATADIR%%/20_fake_helo_tests.cf
%%DATADIR%%/20_head_tests.cf
%%DATADIR%%/20_html_tests.cf
%%DATADIR%%/20_meta_tests.cf
%%DATADIR%%/20_phrases.cf
%%DATADIR%%/20_porn.cf
%%DATADIR%%/20_ratware.cf
%%DATADIR%%/20_uri_tests.cf
%%DATADIR%%/23_bayes.cf
%%DATADIR%%/25_body_tests_es.cf
%%DATADIR%%/25_body_tests_pl.cf
%%DATADIR%%/25_head_tests_es.cf
%%DATADIR%%/25_head_tests_pl.cf
%%DATADIR%%/30_text_de.cf
%%DATADIR%%/30_text_es.cf
%%DATADIR%%/30_text_fr.cf
%%DATADIR%%/30_text_it.cf
%%DATADIR%%/30_text_pl.cf
%%DATADIR%%/30_text_sk.cf
%%DATADIR%%/50_scores.cf
%%DATADIR%%/60_whitelist.cf
%%DATADIR%%/languages
%%DATADIR%%/triplets.txt
%%DATADIR%%/user_prefs.template
@dirrm %%DATADIR%%
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm %%SITE_PERL%%/%%PERL_ARCH%%/auto/Mail/SpamAssassin
@unexec rmdir %D/%%SITE_PERL%%/%%PERL_ARCH%%/auto/Mail 2>/dev/null || true
@dirrm %%SITE_PERL%%/Mail/SpamAssassin
@unexec rmdir %D/%%SITE_PERL%%/Mail 2>/dev/null || true
@unexec rmdir %D/etc/mail/spamassassin 2>/dev/null || true
@unexec rmdir %D/etc/mail 2>/dev/null || true
