# net-im/jitsi-videobridge/files/jitsi-videobridge.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: jitsi-videobridge
# REQUIRE: NETWORKING SERVERS
# KEYWORD: shutdown
#
# Configuration settings for jitsi-videobridge in /etc/rc.conf:
#
# jitsi_videobridge_enable:	run jitsi_videobridge server (default=NO)
# jitsi_videobridge_flags:	additional flags for jitsi_videobridge server (default="")
#

. /etc/rc.subr

name="jitsi_videobridge"
rcvar=${name}_enable

jitsi_videobridge_env_file="%%ETCDIR%%/jitsi-videobridge.conf"
jitsi_videobridge_jar="%%JAVAJARDIR%%/jitsi-videobridge.jar"
jitsi_videobridge_logging_config="%%DATADIR%%/lib/logging.properties"

load_rc_config "${name}"

: ${jitsi_videobridge_enable="NO"}
: ${jitsi_videobridge_flags=""}

command="%%JAVA%%"

pidfile="/var/run/jitsi-videobridge.pid"
required_files="${jitsi_videobridge_env_file}"
start_cmd=${name}_start
stop_cmd=${name}_stop
status_cmd=${name}_status
restart_cmd=${name}_restart

jitsi_videobridge_start()
{
	daemon -p ${pidfile} -o /var/log/jitsi-videobridge.log \
		${command} -Xmx${VIDEOBRIDGE_MAX_MEMORY} \
                -XX:+UseConcMarkSweepGC \
                -XX:+HeapDumpOnOutOfMemoryError \
                -XX:HeapDumpPath=/tmp \
                -Djava.util.logging.config.file=${jitsi_videobridge_logging_config} \
                -Dnet.java.sip.communicator.SC_HOME_DIR_LOCATION=%%LOCALBASE%%/etc/jitsi \
                -Dnet.java.sip.communicator.SC_HOME_DIR_NAME=videobridge \
                -Dnet.java.sip.communicator.SC_LOG_DIR_LOCATION=/var/log/ \
                -cp ${jitsi_videobridge_jar} \
                org.jitsi.videobridge.Main \
                --host=${JVB_XMPP_HOST} \
                --domain=${JVB_XMPP_DOMAIN} \
                --port=${JVB_XMPP_PORT} \
                --secret=${JVB_XMPP_SECRET} ${jitsi_videobridge_flags}
                echo "Started"
}

jitsi_videobridge_stop()
{
        if [ -f ${pidfile} ]
        then
                kill `cat ${pidfile}`
                rm ${pidfile}
		echo "Stopped"
        fi
}

jitsi_videobridge_status()
{
        # If running, show pid
        if [ -f ${pidfile} ]
        then
                echo "${name} is running as pid" `cat ${pidfile}`
        else
                echo "${name} is not running"
        fi
}

jitsi_videobridge_restart()
{
	echo "Performing restart"
	jitsi_videobridge_stop
	sleep 5
	jitsi_videobridge_start
}

run_rc_command "$1"
