# mail/dbmail22/files/openrc-dbmail-timsieved.in
#!/sbin/openrc-run

name=dbmail_timsieved
command=%%PREFIX%%/sbin/dbmail-timsieved
pidfile=/var/run/dbmail-timsieved.pid
required_files=%%PREFIX%%/etc/dbmail.conf

depend() {
	keyword -shutdown
}
