# mail/dbmail22/files/openrc-dbmail-pop3d.in
#!/sbin/openrc-run

name=dbmail_pop3d
command=%%PREFIX%%/sbin/dbmail-pop3d
pidfile=/var/run/dbmail-pop3d.pid
required_files=%%PREFIX%%/etc/dbmail.conf

depend() {
	keyword -shutdown
}
