# multimedia/ffmpeg26/files/ffserver26.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ffserver26
# REQUIRE: NETWORKING
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable ffserver:
#
#ffserver26_enable="YES"
#
ffserver26_enable="${ffserver_enable-NO}"

. /etc/rc.subr

name=ffserver26
rcvar=ffserver26_enable

command="%%PREFIX%%/bin/${name}"
command_args="&"
required_files=%%PREFIX%%/etc/ffserver26.conf

load_rc_config ${name}
run_rc_command "$1"
