# net-mgmt/netustad/pkg-plist
bin/netustadctl
bin/netustadpasswd
sbin/netustad
etc/rc.d/netustad.sh
@unexec if cmp -s %D/etc/netustad/netustad.conf %D/etc/netustad/netustad.conf-dist; then rm -f %D/etc/netustad/netustad.conf; fi
etc/netustad/netustad.conf-dist
@exec [ -f %B/netustad/netustad.conf ] || cp %B/%f %B/netustad/netustad.conf
share/locale/tr/LC_MESSAGES/netustad.mo
%%DATADIR%%/script/action1.inc
%%DATADIR%%/script/action2.inc
%%DATADIR%%/script/edit-parse
%%DATADIR%%/script/if-edit
%%DATADIR%%/script/if-parse
%%DATADIR%%/script/intf-parse
%%DATADIR%%/script/nat-edit
%%DATADIR%%/script/nat-parse
%%DATADIR%%/script/netustad.css
%%DATADIR%%/script/proto.inc
%%DATADIR%%/script/rt-edit
%%DATADIR%%/script/rt-parse
%%DATADIR%%/script/ruleoption1.inc
%%DATADIR%%/script/ruleoption2.inc
%%DATADIR%%/script/service-parse
%%DATADIR%%/script/stats-parse
%%DATADIR%%/script/static-parse
%%DATADIR%%/script/addform
%%DATADIR%%/script/authform
%%DATADIR%%/script/editform
%%DATADIR%%/script/footer
%%DATADIR%%/script/header
%%DATADIR%%/script/if_edit
%%DATADIR%%/script/if_table
%%DATADIR%%/script/loadindex
%%DATADIR%%/script/maketable
%%DATADIR%%/script/nat_add
%%DATADIR%%/script/nat_del
%%DATADIR%%/script/nat_edit
%%DATADIR%%/script/rt_add
%%DATADIR%%/script/rt_del
%%DATADIR%%/script/rt_edit
%%DATADIR%%/script/rt_table
%%DATADIR%%/script/shownat
%%DATADIR%%/script/showrule
%%DATADIR%%/script/stats
%%DATADIR%%/script/addrule.content
%%DATADIR%%/script/editrule.content
%%DATADIR%%/script/writeconfig
%%DATADIR%%/script/images/action_allow.png
%%DATADIR%%/script/images/action_checkstate.png
%%DATADIR%%/script/images/action_count.png
%%DATADIR%%/script/images/action_deny.png
%%DATADIR%%/script/images/action_divert.png
%%DATADIR%%/script/images/action_forward.png
%%DATADIR%%/script/images/action_pipequeue.png
%%DATADIR%%/script/images/action_skipto.png
%%DATADIR%%/script/images/action_tee.png
%%DATADIR%%/script/images/action_unreach.png
%%DATADIR%%/script/images/banner-alt.png
%%DATADIR%%/script/images/banner.png
%%DATADIR%%/script/images/bg.png
%%DATADIR%%/script/images/btn_bg.png
%%DATADIR%%/script/images/buton_0.png
%%DATADIR%%/script/images/buton_1.png
%%DATADIR%%/script/images/buton_2.png
%%DATADIR%%/script/images/buton_3.png
%%DATADIR%%/script/images/buton_4.png
%%DATADIR%%/script/images/buton_5.png
%%DATADIR%%/script/images/buton_6.png
%%DATADIR%%/script/images/buton_7.png
%%DATADIR%%/script/images/editor.png
%%DATADIR%%/script/images/editor_static.png
%%DATADIR%%/script/images/eraser.png
%%DATADIR%%/script/images/eraser_static.png
%%DATADIR%%/script/images/menu_addrule.png
%%DATADIR%%/script/images/menu_login.png
%%DATADIR%%/script/images/menu_logout.png
%%DATADIR%%/script/images/menu_rulelist.png
%%DATADIR%%/cert/netustad.crt
%%DATADIR%%/cert/netustad.csr
%%DATADIR%%/cert/netustad.key
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%%%DOCSDIR%%/ABOUT-NLS
@dirrm etc/netustad
@dirrm share/locale/tr/LC_MESSAGES
@dirrm %%DATADIR%%/script/images
@dirrm %%DATADIR%%/script
@dirrm %%DATADIR%%/cert
@dirrm %%DATADIR%%
%%PORTDOCS%%@dirrm %%DOCSDIR%%
