# databases/postgresql80-server/pkg-plist
etc/rc.d/postgres95.sh
postgres95/lib/libpq.a
postgres95/bin/postgres
postgres95/bin/postmaster
postgres95/bin/pg_id
postgres95/bin/pg_version
postgres95/bin/psql
postgres95/bin/pg_dump
postgres95/bin/cleardbdir
postgres95/bin/createdb
postgres95/bin/createuser
postgres95/bin/destroydb
postgres95/bin/destroyuser
postgres95/bin/initdb
postgres95/data/files/global1.bki
postgres95/data/files/global1.bki.source
postgres95/data/files/local1_template1.bki
postgres95/data/files/local1_template1.bki.source
postgres95/data/pg_hba
postgres95/data/base/template1/pg_proc
postgres95/data/base/template1/pg_type
postgres95/data/base/template1/pg_attribute
postgres95/data/base/template1/pg_class
postgres95/data/base/template1/pg_inherits
postgres95/data/base/template1/pg_index
postgres95/data/base/template1/pg_version
postgres95/data/base/template1/pg_statistic
postgres95/data/base/template1/pg_operator
postgres95/data/base/template1/pg_opclass
postgres95/data/base/template1/pg_am
postgres95/data/base/template1/pg_amop
postgres95/data/base/template1/pg_amproc
postgres95/data/base/template1/pg_language
postgres95/data/base/template1/pg_parg
postgres95/data/base/template1/pg_aggregate
postgres95/data/base/template1/pg_ipl
postgres95/data/base/template1/pg_inheritproc
postgres95/data/base/template1/pg_rewrite
postgres95/data/base/template1/pg_listener
postgres95/data/base/template1/pg_attnameind
postgres95/data/base/template1/pg_attnumind
postgres95/data/base/template1/pg_attrelidind
postgres95/data/base/template1/pg_procidind
postgres95/data/base/template1/pg_procnameind
postgres95/data/base/template1/pg_procsrcind
postgres95/data/base/template1/pg_typeidind
postgres95/data/base/template1/pg_typenameind
postgres95/data/base/template1/pg_classnameind
postgres95/data/base/template1/pg_classoidind
postgres95/data/base/template1/PG_VERSION
postgres95/data/base/template1/pg_internal.init
postgres95/data/pg_variable
postgres95/data/pg_database
postgres95/data/pg_demon
postgres95/data/pg_magic
postgres95/data/pg_defaults
postgres95/data/pg_server
postgres95/data/pg_user
postgres95/data/pg_hosts
postgres95/data/pg_group
postgres95/data/pg_log
postgres95/data/pg_time
postgres95/data/PG_VERSION
postgres95/include/fmgr.h
postgres95/include/lib/dllist.h
postgres95/include/libpq/pqcomm.h
postgres95/include/libpq/libpq-fs.h
postgres95/include/utils/geo-decls.h
postgres95/include/libpq-fe.h
postgres95/include/postgres.h
postgres95/include/c.h
postgres95/man/man1/cleardbdir.1.gz
postgres95/man/man1/psql.1.gz
postgres95/man/man1/postmaster.1.gz
postgres95/man/man1/postgres.1.gz
postgres95/man/man1/pg_dump.1.gz
postgres95/man/man1/monitor.1.gz
postgres95/man/man1/ipcclean.1.gz
postgres95/man/man1/initdb.1.gz
postgres95/man/man1/destroyuser.1.gz
postgres95/man/man1/destroydb.1.gz
postgres95/man/man1/createuser.1.gz
postgres95/man/man1/createdb.1.gz
postgres95/man/man1/unix.1.gz
postgres95/man/man3/libpq.3.gz
postgres95/man/man3/large_objects.3.gz
postgres95/man/man3/catalogs.3.gz
postgres95/man/man3/built-in.3.gz
postgres95/man/man5/template.5.gz
postgres95/man/man5/pg_hba.5.gz
postgres95/man/man5/page.5.gz
postgres95/man/man5/files.5.gz
postgres95/man/man5/bki.5.gz
postgres95/man/manl/abort.l.gz
postgres95/man/manl/update.l.gz
postgres95/man/manl/sql.l.gz
postgres95/man/manl/select.l.gz
postgres95/man/manl/rollback.l.gz
postgres95/man/manl/revoke.l.gz
postgres95/man/manl/rename.l.gz
postgres95/man/manl/remove_view.l.gz
postgres95/man/manl/purge.l.gz
postgres95/man/manl/notify.l.gz
postgres95/man/manl/load.l.gz
postgres95/man/manl/listen.l.gz
postgres95/man/manl/insert.l.gz
postgres95/man/manl/grant.l.gz
postgres95/man/manl/fetch.l.gz
postgres95/man/manl/end.l.gz
postgres95/man/manl/drop_type.l.gz
postgres95/man/manl/drop_rule.l.gz
postgres95/man/manl/drop_operator.l.gz
postgres95/man/manl/drop_index.l.gz
postgres95/man/manl/drop_function.l.gz
postgres95/man/manl/drop_aggregate.l.gz
postgres95/man/manl/drop.l.gz
postgres95/man/manl/destroydb.l.gz
postgres95/man/manl/copy.l.gz
postgres95/man/manl/delete.l.gz
postgres95/man/manl/createdb.l.gz
postgres95/man/manl/create_view.l.gz
postgres95/man/manl/create_version.l.gz
postgres95/man/manl/create_type.l.gz
postgres95/man/manl/create_table.l.gz
postgres95/man/manl/create_rule.l.gz
postgres95/man/manl/create_operator.l.gz
postgres95/man/manl/create_index.l.gz
postgres95/man/manl/create_function.l.gz
postgres95/man/manl/create_aggregate.l.gz
postgres95/man/manl/commit.l.gz
postgres95/man/manl/cluster.l.gz
postgres95/man/manl/close.l.gz
postgres95/man/manl/begin.l.gz
postgres95/man/manl/alter_table.l.gz
postgres95/man/manl/vacuum.l.gz
postgres95/post-install-notes
@dirrm postgres95/man/man1/
@dirrm postgres95/man/man3/
@dirrm postgres95/man/man5/
@dirrm postgres95/man/manl/
@dirrm postgres95/man/
@dirrm postgres95/lib/
@dirrm postgres95/include/utils/
@dirrm postgres95/include/port/BSD44_derived/
@dirrm postgres95/include/port/
@dirrm postgres95/include/libpq/
@dirrm postgres95/include/lib/
@dirrm postgres95/include/
@dirrm postgres95/data/files/
@dirrm postgres95/data/base/template1/
@dirrm postgres95/data/base/
@dirrm postgres95/data/
@dirrm postgres95/bin/
@dirrm postgres95/
