# dns/dynip/pkg-plist
bin/dynip
etc/rc.d/dynip
etc/periodic/daily/dynip
share/man/man8/dynip.8.gz
