# security/tor-devel/pkg-plist
@comment $FreeBSD$
bin/tor
bin/tor-resolve
bin/torify
etc/rc.d/tor.sh.sample
etc/tor/torrc.sample
etc/tor/tor-tsocks.conf.sample
@unexec rmdir %D/etc/tor 2>/dev/null || true
