# mail/sccmilter/pkg-plist
etc/rc.d/sccmilter.sh.sample
etc/sccmilter/scchosts.sample
etc/sccmilter/sccmilter.conf.sample
etc/sccmilter/sccusers.sample
libexec/sccmilter
%%PORTDOCS%%%%DOCSDIR%%/CHANGES
%%PORTDOCS%%%%DOCSDIR%%/FAQ
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm etc/sccmilter
