# audio/oss/files/openrc-oss.in
#!/sbin/openrc-run

name="oss"

start()
{
	ebegin "Starting Open Sound System"
	command=%%PREFIX%%/sbin/soundon
	$command
}

stop()
{
	eend "Stopping Open Sound System"
	command=%%PREFIX%%/sbin/soundoff
	$command
}

depend() {
	provide oss
	keyword -shutdown -jail
}
