# net-mgmt/zabbix24-server/files/zabbix_proxy.in
#!/bin/sh

# PROVIDE: zabbix_proxy
# REQUIRE: DAEMON
%%PGSQL%%# REQUIRE: postgresql
%%MYSQL%%# REQUIRE: mysql-server
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf to
# enable zabbix_proxy:
#
# zabbix_proxy_enable (bool): Set to NO by default.  Set it to YES to
#         enable zabbix_proxy.
# zabbix_proxy_paths (string): Set to standard path by default.  Set a search
#         if you have custom externals that need binaries elsewhere.
# zabbix_proxy_config (string): Set to the standard config file path by
#         default.
#

. /etc/rc.subr

name="zabbix_proxy"
rcvar=zabbix_proxy_enable

load_rc_config $name

: ${zabbix_proxy_enable="NO"}
: ${zabbix_proxy_paths=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/sbin:/usr/local/bin}
: ${zabbix_proxy_config="%%ETCDIR%%/${name}.conf"}

command="%%PREFIX%%/sbin/${name}"
required_files="${zabbix_proxy_config}"
start_precmd="find_pidfile"
status_precmd="find_pidfile"
stop_precmd="find_pidfile"

find_pidfile()
{
	if get_pidfile_from_conf PidFile ${zabbix_agentd_config}; then
		pidfile="$_pidfile_from_conf"
	else
		pidfile="/tmp/${name}.pid"
	fi

	# This shouldn't be necessary with pidfile, but empirically it was the
	# only way to reap the parent PID instead of all PIDs from
	# check_process, which may leak SysV IPC objects and prevent restart
	# and/or race condition on restart.
	rc_pid=$(check_pidfile ${pidfile} ${command})
}

export PATH="${zabbix_proxy_paths}"

run_rc_command "$1"
