# mail/cyrus-imapd22/pkg-plist
bin/cyradm
bin/imtest
cyrus/bin/pwcheck
cyrus/bin/arbitron
cyrus/bin/collectnews
cyrus/bin/cyrquota
cyrus/bin/deliver
cyrus/bin/dump_deliverdb
cyrus/bin/feedcyrus
cyrus/bin/fud
cyrus/bin/imapd
cyrus/bin/pop3d
cyrus/bin/reconstruct
cyrus/bin/syncnews
cyrus/bin/updateimsp
etc/imapd.conf
etc/rc.d/cyrus.sh
include/cyrus/acl.h
include/cyrus/assert.h
include/cyrus/auth.h
include/cyrus/bsearch.h
include/cyrus/charset.h
include/cyrus/glob.h
include/cyrus/gmtoff.h
include/cyrus/imclient.h
include/cyrus/imparse.h
include/cyrus/lock.h
include/cyrus/map.h
include/cyrus/mkgmtime.h
include/cyrus/nonblock.h
include/cyrus/parseaddr.h
include/cyrus/prot.h
include/cyrus/retry.h
include/cyrus/sasl.h
include/cyrus/sysexits.h
include/cyrus/util.h
include/cyrus/xmalloc.h
lib/libcyrus.a
share/doc/cyrus/html/bugs.html
share/doc/cyrus/html/changes.html
share/doc/cyrus/html/htmlstrip.c
share/doc/cyrus/html/index.html
share/doc/cyrus/html/install.html
share/doc/cyrus/html/overview.html
share/doc/cyrus/html/readme.html
share/doc/cyrus/html/unpack.html
share/doc/cyrus/acl-extension
share/doc/cyrus/bugs
share/doc/cyrus/changes
share/doc/cyrus/copyrights
share/doc/cyrus/install
share/doc/cyrus/overview
share/doc/cyrus/quota-extension
share/doc/cyrus/server-design
@dirrm include/cyrus
@dirrm cyrus/bin
@dirrm cyrus
@dirrm share/doc/cyrus/html
@dirrm share/doc/cyrus
