# sysutils/anacron/pkg-plist
sbin/anacron
etc/anacrontab.sample
etc/rc.d/anacron.sh
