# mail/assp/files/openrc-assp.in
#!/sbin/openrc-run

name=assp
command=%%PREFIX%%/sbin/assp
pidfile=%%ASSP_HOME%%/pid
: ${assp_flags="%%ASSP_HOME%%"}
command_args=${assp_flags}

depend() {
	keyword -shutdown
}
