# net/pim6sd/files/openrc-pim6sd.in
#!/sbin/openrc-run

name="pim6sd"
command="%%PREFIX%%/sbin/${name}"
required_files="%%PREFIX%%/etc/${name}.conf"

depend() {
	keyword -jail
}
