# sysutils/cfengine22/files/cfservd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: cfservd
# REQUIRE: SERVERS
# BEFORE:  DAEMON
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf[.local] to enable cfservd(8)
#
# cfservd_enable (bool):        Set to "NO" by default.
#                               Set it to "YES" to enable cfservd.
# cfservd_flags (str):          Custom additional arguments to be passed
#                               to cfservd (default empty).
#

. /etc/rc.subr

name="cfservd"
rcvar=cfservd_enable

command="%%PREFIX%%/sbin/${name}"

load_rc_config $name

: ${cfservd_enable="NO"}

run_rc_command "$1"
