# sysutils/torque/files/openrc-pbs_server.in
#!/sbin/openrc-run

name=pbs_server

pidfile=/var/spool/torque/server_priv/server.lock
command=%%PREFIX%%/sbin/${name}

depend() {
	keyword -shutdown
}
