# sysutils/torque/files/openrc-pbs_sched.in
#!/sbin/openrc-run

name=pbs_sched

pidfile=/var/spool/torque/sched_priv/sched.lock
command=%%PREFIX%%/sbin/${name}

depend() {
	keyword -shutdown
}
