# sysutils/torque/files/pbs_mom.sh.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: pbs_mom
# REQUIRE: 
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# pbs_mom_enable (bool):   Set to NO by default.
#               Set it to YES to enable pbs_mom.
#

. /etc/rc.subr

name="pbs_mom"
rcvar=pbs_mom_enable

pidfile=/var/spool/torque/mom_priv/mom.lock
command=%%PREFIX%%/sbin/${name}

load_rc_config $name

: ${pbs_mom_enable="NO"}

run_rc_command "$1"
