# sysutils/torque/files/pbs_server.sh.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: pbs_server
# REQUIRE: 
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# pbs_server_enable (bool):   Set to NO by default.
#               Set it to YES to enable pbs_server.
#

. /etc/rc.subr

name="pbs_server"
rcvar=pbs_server_enable

pidfile=/var/spool/torque/server_priv/server.lock
command=%%PREFIX%%/sbin/${name}

load_rc_config $name

: ${pbs_server_enable="NO"}

run_rc_command "$1"
