# databases/postgresql82-server/pkg-plist
@unexec if cmp -s %D/pgsql/.profile %D/pgsql/dot.profile.dist; then rm -f %D/pgsql/.profile; fi
etc/rc.d/010.pgsql.sh
pgsql/dot.profile.dist
@exec [ -f %D/pgsql/.profile ] || (cp %D/pgsql/dot.profile.dist %D/pgsql/.profile && chown pgsql:pgsql %D/pgsql/.profile)
pgsql/bin/createdb
pgsql/bin/createlang
pgsql/bin/createuser
pgsql/bin/dropdb
pgsql/bin/droplang
pgsql/bin/dropuser
pgsql/bin/ecpg
pgsql/bin/initdb
pgsql/bin/initlocation
pgsql/bin/ipcclean
pgsql/bin/pg_ctl
pgsql/bin/pg_dump
pgsql/bin/pg_dumpall
pgsql/bin/pg_encoding
pgsql/bin/pg_id
pgsql/bin/pg_passwd
pgsql/bin/pg_upgrade
pgsql/bin/pg_version
pgsql/bin/postgres
pgsql/bin/postmaster
pgsql/bin/psql
pgsql/bin/vacuumdb
pgsql/include/access/attnum.h
pgsql/include/c.h
pgsql/include/commands/trigger.h
pgsql/include/config.h
pgsql/include/ecpgerrno.h
pgsql/include/ecpglib.h
pgsql/include/ecpgtype.h
pgsql/include/executor/spi.h
pgsql/include/fmgr.h
pgsql/include/lib/dllist.h
pgsql/include/libpgeasy.h
pgsql/include/libpq/libpq-fs.h
pgsql/include/libpq/pqcomm.h
pgsql/include/libpq++/pgconnection.h
pgsql/include/libpq++/pgcursordb.h
pgsql/include/libpq++/pgdatabase.h
pgsql/include/libpq++/pglobject.h
pgsql/include/libpq++/pgtransdb.h
pgsql/include/libpq++.h
pgsql/include/libpq-fe.h
pgsql/include/libpq-int.h
pgsql/include/os.h
pgsql/include/postgres.h
pgsql/include/postgres_ext.h
pgsql/include/pqexpbuffer.h
pgsql/include/sql3types.h
pgsql/include/sqlca.h
pgsql/include/utils/elog.h
pgsql/include/utils/geo_decls.h
pgsql/include/utils/mcxt.h
pgsql/include/utils/palloc.h
pgsql/lib/global1.bki.source
pgsql/lib/global1.description
pgsql/lib/libecpg.a
pgsql/lib/libecpg.so
pgsql/lib/libecpg.so.3
pgsql/lib/libpgeasy.a
pgsql/lib/libpgeasy.so
pgsql/lib/libpgeasy.so.2
pgsql/lib/libpq++.a
pgsql/lib/libpq++.so
pgsql/lib/libpq++.so.3
pgsql/lib/libpq.a
pgsql/lib/libpq.so
pgsql/lib/libpq.so.2
pgsql/lib/local1_template1.bki.source
pgsql/lib/local1_template1.description
pgsql/lib/pg_geqo.sample
pgsql/lib/pg_hba.conf.sample
pgsql/lib/pg_ident.conf.sample
pgsql/lib/pg_options.sample
pgsql/lib/plpgsql.so
pgsql/lib/postmaster.opts.default.sample
@exec chown -R root:wheel %D/pgsql/lib
pgsql/post-install-notes
@dirrm pgsql/bin
@dirrm pgsql/include/access
@dirrm pgsql/include/commands
@dirrm pgsql/include/executor
@dirrm pgsql/include/lib
@dirrm pgsql/include/libpq
@dirrm pgsql/include/libpq++
@dirrm pgsql/include/utils
@dirrm pgsql/include
@dirrm pgsql/lib
@dirrm pgsql/man/man1
@dirrm pgsql/man/manl
@dirrm pgsql/man
@dirrm pgsql
