# net/quagga/files/openrc-quagga.in
#!/sbin/openrc-run

#  quagga_flags=""
#  quagga_daemons="zebra ripd ripngd ospfd ospf6d bgpd isisd"
#  quagga_vysh_boot="YES"
# Per daemon tuning may be done with daemon_name_flags
#  zebra_flags="-P 0"
#  bgpd_flags="-nrP 0" and so on
# If you want to give the routing deamons a chance to catchup before
# continueing, set quagga_wait_for to a "default" or certain prefix.
#  quagga_wait_for="default"
# Set the time limit for the wait.
#  quagga_wait_seconds="90"
#
# If the quagga daemons require additional shared libraries to start,
# use the following variable to run ldconfig(8) in advance:
#  quagga_extralibs_path="/usr/local/lib ..."
#

name=quagga

extra_commands=configtest
command_args="-d"
quagga_flags=${quagga_flags:-""}
quagga_daemons=${quagga_daemons:-"zebra ripd ripngd ospfd ospf6d bgpd isisd"}
quagga_vtysh_boot=${quagga_vtysh_boot:-"NO"}
quagga_wait_for=${quagga_wait_for:-""}
quagga_wait_seconds=${quagga_wait_seconds:-"90"}

configtest()
{
    einfo "Checking $daemon.conf"
    $command $daemon_flags -C
    result=$?
    if [ "$result" -eq "0" ]; then
		einfo "OK"
    else
		eerror "FAILED"
		exit
    fi
}

start_post()
{
	local waited_for
	waited_for=0
	# Wait only when last daemon has started.
	if [ "${quagga_daemons}" = "${quagga_daemons% ${name}}" ]; then
		return;
	fi
        if [ -n "${quagga_wait_for}" ]; then
		einfo Waiting for ${quagga_wait_for} route...
		while [ ${waited_for} -lt ${quagga_wait_seconds} ]; do
			/sbin/route -n get ${quagga_wait_for} >/dev/null 2>&1 && break;
			waited_for=$((${waited_for}+1))
			sleep 1;
		done
		[ ${waited_for} -lt ${quagga_wait_seconds} ] || einfo Giving up...
	fi
}

start_pre() {
	if [ $# -ge 1 ]; then
		quagga_daemons="$*"
	fi
	if [ -n "${quagga_extralibs_path}" ]; then
	    /sbin/ldconfig -m ${quagga_extralibs_path}
	fi
	check_config
	for daemon in ${quagga_daemons}; do
	    command=%%PREFIX%%/sbin/${daemon}
	    required_files=%%SYSCONF_DIR%%/${daemon}.conf
	    pidfile=%%LOCALSTATE_DIR%%/${daemon}.pid
	    eval flags=\$\{${daemon}_flags:-\"${quagga_flags}\"\}
	    name=${daemon}
	done
	if yesno quagga_vtysh_boot; then
	    einfo "Booting for integrated-vtysh-config..."
	    %%PREFIX%%/bin/vtysh -b
	fi
}

stop_pre() {
	quagga_daemons=$(reverse_list ${quagga_daemons})
}

depend() {
	provide quagga
	keyword -shutdown -jail
}
