# net/openldap22-server/pkg-plist
@comment $FreeBSD$
@unexec %%RC_DIR%%/etc/rc.d/slapd%%RC_SUFX%% stop 2>&1 >/dev/null || true
@unexec %%RC_DIR%%/etc/rc.d/slurpd%%RC_SUFX%% stop 2>&1 >/dev/null || true
etc/openldap/schema/README
@unexec if cmp -s %D/etc/openldap/slapd.conf %D/etc/openldap/slapd.conf.default; then rm -f %D/etc/openldap/slapd.conf; fi
etc/openldap/slapd.conf.default
@exec [ -f %B/slapd.conf ] || cp %B/%f %B/slapd.conf
%%SLAPI%%lib/libslapi.a
%%SLAPI%%lib/libslapi.so
%%SLAPI%%lib/libslapi-2.2.so
%%SLAPI%%lib/libslapi-2.2.so.7
%%MODULES%%@exec mkdir -p %D/libexec/openldap
%%BACK_BDB%%libexec/openldap/back_bdb.so
%%BACK_BDB%%libexec/openldap/back_bdb-2.2.so
%%BACK_BDB%%libexec/openldap/back_bdb-2.2.so.7
%%BACK_HDB%%libexec/openldap/back_hdb.so
%%BACK_HDB%%libexec/openldap/back_hdb-2.2.so
%%BACK_HDB%%libexec/openldap/back_hdb-2.2.so.7
%%BACKEND%%libexec/openldap/back_ldap.so
%%BACKEND%%libexec/openldap/back_ldap-2.2.so
%%BACKEND%%libexec/openldap/back_ldap-2.2.so.7
%%BACKEND%%libexec/openldap/back_ldbm.so
%%BACKEND%%libexec/openldap/back_ldbm-2.2.so
%%BACKEND%%libexec/openldap/back_ldbm-2.2.so.7
%%BACKEND%%libexec/openldap/back_meta.so
%%BACKEND%%libexec/openldap/back_meta-2.2.so
%%BACKEND%%libexec/openldap/back_meta-2.2.so.7
%%BACKEND%%libexec/openldap/back_monitor.so
%%BACKEND%%libexec/openldap/back_monitor-2.2.so
%%BACKEND%%libexec/openldap/back_monitor-2.2.so.7
%%BACKEND%%libexec/openldap/back_null.so
%%BACKEND%%libexec/openldap/back_null-2.2.so
%%BACKEND%%libexec/openldap/back_null-2.2.so.7
%%BACK_PERL%%libexec/openldap/back_perl.so
%%BACK_PERL%%libexec/openldap/back_perl-2.2.so
%%BACK_PERL%%libexec/openldap/back_perl-2.2.so.7
%%BACK_SHELL%%libexec/openldap/back_shell.so
%%BACK_SHELL%%libexec/openldap/back_shell-2.2.so
%%BACK_SHELL%%libexec/openldap/back_shell-2.2.so.7
%%BACK_SQL%%libexec/openldap/back_sql.so
%%BACK_SQL%%libexec/openldap/back_sql-2.2.so
%%BACK_SQL%%libexec/openldap/back_sql-2.2.so.7
%%MODULES%%@unexec rmdir %D/libexec/openldap 2>/dev/null || true
libexec/slapd
libexec/slurpd
@comment sbin/slapacl
sbin/slapadd
@comment sbin/slapauth
sbin/slapcat
sbin/slapdn
sbin/slapindex
sbin/slappasswd
sbin/slaptest
@exec mkdir -p %%LDAP_RUN_DIR%%
@unexec rmdir %%LDAP_RUN_DIR%% 2>/dev/null || true
@exec mkdir -p %%DATABASEDIR%%
@unexec rmdir %%DATABASEDIR%% 2>/dev/null || true
@exec mkdir -p %%SLURPDIR%%
@unexec rmdir %%SLURPDIR%% 2>/dev/null || true
