# sysutils/appjail-devel/pkg-plist
bin/appjail
bin/appjail-config
bin/appjail-config-user
bin/appjail-dns
bin/appjail-user
@sample %%ETCDIR%%/appjail.conf.sample
etc/bash_completion.d/_appjail.bash
etc/rc.d/appjail
etc/rc.d/appjail-dns
etc/rc.d/appjail-health
etc/rc.d/appjail-natnet
libexec/appjail/appjail-config/appjail-config
libexec/appjail/appjail-config/tok
libexec/appjail/find-number-from-start/find-number-from-start
libexec/appjail/find-smallest-missing-number/find-smallest-missing-number
libexec/appjail/get_assigned_rulesets/get_assigned_rulesets
libexec/appjail/getservbyname/getservbyname
libexec/appjail/ipcheck/ipcheck
libexec/appjail/jexec/jexec
libexec/appjail/network/network
%%DATADIR%%/cmd/apply
%%DATADIR%%/cmd/checkOld
%%DATADIR%%/cmd/cmd
%%DATADIR%%/cmd/cpuset
%%DATADIR%%/cmd/deleteOld
%%DATADIR%%/cmd/devfs
%%DATADIR%%/cmd/disable
%%DATADIR%%/cmd/enable
%%DATADIR%%/cmd/enabled
%%DATADIR%%/cmd/etcupdate
%%DATADIR%%/cmd/expose
%%DATADIR%%/cmd/fetch
%%DATADIR%%/cmd/fstab
%%DATADIR%%/cmd/healthcheck
%%DATADIR%%/cmd/help
%%DATADIR%%/cmd/image
%%DATADIR%%/cmd/jail
%%DATADIR%%/cmd/label
%%DATADIR%%/cmd/limits
%%DATADIR%%/cmd/login
%%DATADIR%%/cmd/logs
%%DATADIR%%/cmd/makejail
%%DATADIR%%/cmd/nat
%%DATADIR%%/cmd/network
%%DATADIR%%/cmd/oci
%%DATADIR%%/cmd/pkg
%%DATADIR%%/cmd/quick
%%DATADIR%%/cmd/restart
%%DATADIR%%/cmd/rstop
%%DATADIR%%/cmd/run
%%DATADIR%%/cmd/service
%%DATADIR%%/cmd/start
%%DATADIR%%/cmd/startup
%%DATADIR%%/cmd/status
%%DATADIR%%/cmd/stop
%%DATADIR%%/cmd/sysrc
%%DATADIR%%/cmd/update
%%DATADIR%%/cmd/upgrade
%%DATADIR%%/cmd/usage
%%DATADIR%%/cmd/version
%%DATADIR%%/cmd/volume
%%DATADIR%%/cmd/zfs
%%DATADIR%%/files/config.conf
%%DATADIR%%/files/default.conf
%%DATADIR%%/files/default_template.conf
%%DATADIR%%/files/dnsmasq.conf
%%DATADIR%%/files/exclude-thinjail.files
%%DATADIR%%/files/include-thinjail.files
%%DATADIR%%/files/rctl-actions
%%DATADIR%%/files/rctl-resources
%%DATADIR%%/files/rctl-subjects
%%DATADIR%%/files/pkgbase/base.conf
%%DATADIR%%/lib/atexit
%%DATADIR%%/lib/check_func
%%DATADIR%%/lib/colors
%%DATADIR%%/lib/copy
%%DATADIR%%/lib/devfs
%%DATADIR%%/lib/files
%%DATADIR%%/lib/initscript
%%DATADIR%%/lib/jail
%%DATADIR%%/lib/jail_types
%%DATADIR%%/lib/kern_modules
%%DATADIR%%/lib/keys
%%DATADIR%%/lib/load
%%DATADIR%%/lib/log
%%DATADIR%%/lib/mksum
%%DATADIR%%/lib/mount
%%DATADIR%%/lib/network
%%DATADIR%%/lib/oci
%%DATADIR%%/lib/random
%%DATADIR%%/lib/repeat
%%DATADIR%%/lib/replace
%%DATADIR%%/lib/select
%%DATADIR%%/lib/strlen
%%DATADIR%%/lib/sysexits
%%DATADIR%%/lib/table
%%DATADIR%%/lib/tempfile
%%DATADIR%%/lib/version
%%DATADIR%%/lib/whitespaces
%%DATADIR%%/lib/zfs
%%DATADIR%%/makejail/cmd/all/ADD
%%DATADIR%%/makejail/cmd/all/ARG
%%DATADIR%%/makejail/cmd/all/CLEAR
%%DATADIR%%/makejail/cmd/all/CMD
%%DATADIR%%/makejail/cmd/all/COPY
%%DATADIR%%/makejail/cmd/all/DESTROY
%%DATADIR%%/makejail/cmd/all/DEVICE
%%DATADIR%%/makejail/cmd/all/ENTRYPOINT
%%DATADIR%%/makejail/cmd/all/ENV
%%DATADIR%%/makejail/cmd/all/EXEC
%%DATADIR%%/makejail/cmd/all/GLOBAL
%%DATADIR%%/makejail/cmd/all/LABEL
%%DATADIR%%/makejail/cmd/all/MOUNT
%%DATADIR%%/makejail/cmd/all/PKG
%%DATADIR%%/makejail/cmd/all/RAW
%%DATADIR%%/makejail/cmd/all/REPLACE
%%DATADIR%%/makejail/cmd/all/RUN
%%DATADIR%%/makejail/cmd/all/SERVICE
%%DATADIR%%/makejail/cmd/all/SET
%%DATADIR%%/makejail/cmd/all/SYSRC
%%DATADIR%%/makejail/cmd/all/UMOUNT
%%DATADIR%%/makejail/cmd/all/USER
%%DATADIR%%/makejail/cmd/all/VAR
%%DATADIR%%/makejail/cmd/all/VOLUME
%%DATADIR%%/makejail/cmd/all/WORKDIR
%%DATADIR%%/makejail/cmd/build/FROM
%%DATADIR%%/makejail/cmd/build/LOGIN
%%DATADIR%%/makejail/cmd/build/OPTION
%%DATADIR%%/makejail/cmd/build/RESTART
%%DATADIR%%/makejail/cmd/build/START
%%DATADIR%%/makejail/cmd/build/STOP
%%DATADIR%%/makejail/write/all/ADD
%%DATADIR%%/makejail/write/all/ARG
%%DATADIR%%/makejail/write/all/CMD
%%DATADIR%%/makejail/write/all/COPY
%%DATADIR%%/makejail/write/all/DESTROY
%%DATADIR%%/makejail/write/all/DEVICE
%%DATADIR%%/makejail/write/all/EXEC
%%DATADIR%%/makejail/write/all/GLOBAL
%%DATADIR%%/makejail/write/all/LABEL
%%DATADIR%%/makejail/write/all/MOUNT
%%DATADIR%%/makejail/write/all/PKG
%%DATADIR%%/makejail/write/all/RAW
%%DATADIR%%/makejail/write/all/REPLACE
%%DATADIR%%/makejail/write/all/RUN
%%DATADIR%%/makejail/write/all/SERVICE
%%DATADIR%%/makejail/write/all/SET
%%DATADIR%%/makejail/write/all/SYSRC
%%DATADIR%%/makejail/write/all/UMOUNT
%%DATADIR%%/makejail/write/all/VAR
%%DATADIR%%/makejail/write/all/VOLUME
%%DATADIR%%/makejail/write/all/WORKDIR
%%DATADIR%%/makejail/write/build/FROM
%%DATADIR%%/makejail/write/build/LOGIN
%%DATADIR%%/makejail/write/build/OPTION
%%DATADIR%%/makejail/write/build/RESTART
%%DATADIR%%/makejail/write/build/START
%%DATADIR%%/makejail/write/build/STOP
%%DATADIR%%/scripts/ajconf-user.sh
%%DATADIR%%/scripts/ajconf.sh
%%DATADIR%%/scripts/ajdns.sh
%%DATADIR%%/scripts/ajuser.sh
%%DATADIR%%/scripts/appjail-completion.bash
%%DATADIR%%/scripts/ascii2oct.sh
%%DATADIR%%/scripts/bind-ddns-hook.sh
%%DATADIR%%/scripts/create-bridge.sh
%%DATADIR%%/scripts/dns.sh
%%DATADIR%%/scripts/dnsmasq-hook.sh
%%DATADIR%%/scripts/dummy-hook.sh
%%DATADIR%%/scripts/escape-env-val.sh
%%DATADIR%%/scripts/get_depends.sh
%%DATADIR%%/scripts/git.sh
%%DATADIR%%/scripts/kill_child.sh
%%DATADIR%%/scripts/kill_tree.sh
%%DATADIR%%/scripts/load-devfs-rules.sh
%%DATADIR%%/scripts/run_init.sh
%%DATADIR%%/scripts/run_jail.sh
%%DATADIR%%/scripts/runas.sh
%%DATADIR%%/scripts/super-mv.sh
%%EXAMPLESDIR%%/appjail.conf
%%EXAMPLESDIR%%/templates/bridge.conf
%%EXAMPLESDIR%%/templates/freebsd-oci.conf
%%EXAMPLESDIR%%/templates/linux.conf
%%EXAMPLESDIR%%/templates/vnet.conf
%%EXAMPLESDIR%%/templates/xorg.conf
share/man/man1/appjail-apply.1.gz
share/man/man1/appjail-checkOld.1.gz
share/man/man1/appjail-cmd.1.gz
share/man/man1/appjail-config.1.gz
share/man/man1/appjail-cpuset.1.gz
share/man/man1/appjail-deleteOld.1.gz
share/man/man1/appjail-devfs.1.gz
share/man/man1/appjail-disable.1.gz
share/man/man1/appjail-enable.1.gz
share/man/man1/appjail-enabled.1.gz
share/man/man1/appjail-etcupdate.1.gz
share/man/man1/appjail-expose.1.gz
share/man/man1/appjail-fetch.1.gz
share/man/man1/appjail-fstab.1.gz
share/man/man1/appjail-healthcheck.1.gz
share/man/man1/appjail-help.1.gz
share/man/man1/appjail-image.1.gz
share/man/man1/appjail-jail.1.gz
share/man/man1/appjail-label.1.gz
share/man/man1/appjail-limits.1.gz
share/man/man1/appjail-login.1.gz
share/man/man1/appjail-logs.1.gz
share/man/man1/appjail-makejail.1.gz
share/man/man1/appjail-nat.1.gz
share/man/man1/appjail-network.1.gz
share/man/man1/appjail-oci.1.gz
share/man/man1/appjail-pkg.1.gz
share/man/man1/appjail-quick.1.gz
share/man/man1/appjail-restart.1.gz
share/man/man1/appjail-rstop.1.gz
share/man/man1/appjail-run.1.gz
share/man/man1/appjail-service.1.gz
share/man/man1/appjail-start.1.gz
share/man/man1/appjail-startup.1.gz
share/man/man1/appjail-status.1.gz
share/man/man1/appjail-stop.1.gz
share/man/man1/appjail-sysrc.1.gz
share/man/man1/appjail-update.1.gz
share/man/man1/appjail-upgrade.1.gz
share/man/man1/appjail-usage.1.gz
share/man/man1/appjail-version.1.gz
share/man/man1/appjail-volume.1.gz
share/man/man1/appjail-zfs.1.gz
share/man/man1/appjail.1.gz
share/man/man5/appjail-ajspec.5.gz
share/man/man5/appjail-initscript.5.gz
share/man/man5/appjail-makejail.5.gz
share/man/man5/appjail-template.5.gz
share/man/man5/appjail.conf.5.gz
share/man/man7/appjail-ephemeral.7.gz
share/man/man7/appjail-tutorial.7.gz
share/man/man8/appjail-dns.8.gz
share/man/man8/appjail-user.8.gz
