# security/nessus-devel/pkg-plist
etc/rc.d/nessusd.sh
bin/nessus
bin/nessus-mkcert-client
bin/nessus-mkrand
include/nessus/config.h
include/nessus/includes.h
include/nessus/nessus-devel.h
include/nessus/nessusicmp.h
include/nessus/nessusip.h
include/nessus/nessusraw.h
include/nessus/nessustcp.h
include/nessus/nessusudp.h
include/nessus/ntcompat.h
sbin/nessus-adduser
sbin/nessus-check-signature
sbin/nessus-mkcert
sbin/nessus-rmuser
sbin/nessusd
@exec mkdir -p %D/var/nessus/jobs
@exec mkdir -p %D/var/nessus/logs
@exec mkdir -p %D/var/nessus/users
@unexec rmdir %D/var/nessus/jobs 2>/dev/null || true
@unexec rmdir %D/var/nessus/logs 2>/dev/null || true
@unexec rmdir %D/var/nessus/users 2>/dev/null || true
@unexec rmdir %D/var/nessus 2>/dev/null || true
@unexec rmdir %D/etc/nessus 2>/dev/null || true
@exec mkdir -p %D/lib/nessus/plugins
@unexec rmdir %D/lib/nessus/plugins 2>/dev/null || true
@unexec rmdir %D/lib/nessus 2>/dev/null || true
