# shells/wapsh/pkg-plist
sbin/htshd
etc/rc.d/htshd.sh
@unexec if cmp -s %D/etc/wapsh/htsh_profile %D/etc/wapsh/htsh_profile.dist; then rm -f %D/etc/wapsh/htsh_profile; fi
etc/wapsh/htsh_profile.dist
@exec [ -f %B/htsh_profile ] || cp %B/%f %B/htsh_profile
@unexec if cmp -s %D/etc/wapsh/htshd.conf %D/etc/wapsh/htshd.conf.dist; then rm -f %D/etc/wapsh/htshd.conf; fi
etc/wapsh/htshd.conf.dist
@exec [ -f %B/htshd.conf ] || cp %B/%f %B/htshd.conf
@unexec if cmp -s %D/etc/wapsh/loginhosts %D/etc/wapsh/loginhosts.dist; then rm -f %D/etc/wapsh/loginhosts; fi
etc/wapsh/loginhosts.dist
@exec [ -f %B/loginhosts ] || cp %B/%f %B/loginhosts
%%DATADIR%%/htdocs/command.php
%%DATADIR%%/htdocs/command_html.inc
%%DATADIR%%/htdocs/command_wml.inc
%%DATADIR%%/htdocs/common.inc
%%DATADIR%%/htdocs/common_html.inc
%%DATADIR%%/htdocs/common_wml.inc
%%DATADIR%%/htdocs/controls_wml.php
%%DATADIR%%/htdocs/displayfuncs.inc
%%DATADIR%%/htdocs/getloginhosts.inc
%%DATADIR%%/htdocs/history_edit_wml.php
%%DATADIR%%/htdocs/history_html.php
%%DATADIR%%/htdocs/history_wml.php
%%DATADIR%%/htdocs/login_html.php
%%DATADIR%%/htdocs/login_wml.php
%%DATADIR%%/htdocs/scroll_wml.php
%%DATADIR%%/htdocs/shortcuts_wml.php
%%PORTDOCS%%%%DOCSDIR%%/Image1.jpg
%%PORTDOCS%%%%DOCSDIR%%/Image2.jpg
%%PORTDOCS%%%%DOCSDIR%%/Image3.jpg
%%PORTDOCS%%%%DOCSDIR%%/admin_guide/index.html
%%PORTDOCS%%%%DOCSDIR%%/implementation_notes/image57.gif
%%PORTDOCS%%%%DOCSDIR%%/implementation_notes/index.html
%%PORTDOCS%%%%DOCSDIR%%/index.html
%%PORTDOCS%%%%DOCSDIR%%/qpl.html
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image10.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image11.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image12.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image13.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image14.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image15.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image16.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image17.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image18.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image19.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image20.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image21.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image22.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image23.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image24.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image25.gif
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image26.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image27.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image28.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image29.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image30.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image31.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image32.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image33.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image34.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image35.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image36.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image37.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image38.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image39.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image4.gif
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image40.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image41.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image42.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image43.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image44.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image45.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image46.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image47.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image48.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image49.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image5.gif
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image50.gif
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image51.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image52.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image53.gif
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image54.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image55.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image56.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image6.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image7.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image8.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/Image9.jpg
%%PORTDOCS%%%%DOCSDIR%%/user_guide/index.html
@exec echo "updating /etc/services"; cp /etc/services /etc/services.bak; (grep -v ^htsh /etc/services.bak ; echo "htsh		3001/tcp    # htsh/wapsh server") > /etc/services
@unexec echo "updating /etc/services"; cp /etc/services /etc/services.bak; (grep -v ^htsh /etc/services.bak ) > /etc/services
@dirrm etc/wapsh
%%PORTDOCS%%@dirrm %%DOCSDIR%%/admin_guide
%%PORTDOCS%%@dirrm %%DOCSDIR%%/implementation_notes
%%PORTDOCS%%@dirrm %%DOCSDIR%%/user_guide
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm %%DATADIR%%/htdocs
@dirrm %%DATADIR%%
