# irc/dancer-ircd/pkg-plist
sbin/dancer-ircd
@unexec if cmp -s %D/etc/dancer-ircd/ircd.conf %D/etc/dancer-ircd/ircd.conf.sample; then rm -f %D/etc/dancer-ircd/ircd.conf; fi
etc/dancer-ircd/ircd.conf.sample
@exec if [ ! -f %B/ircd.conf ]; then cp %B/%f %B/ircd.conf; fi
@unexec if cmp -s %D/etc/dancer-ircd/motd %D/etc/dancer-ircd/motd.sample; then rm -f %D/etc/dancer-ircd/motd; fi
etc/dancer-ircd/motd.sample
@exec if [ ! -f %B/motd ]; then cp %B/%f %B/motd; fi
@unexec if cmp -s %D/etc/dancer-ircd/kline.conf %D/etc/dancer-ircd/kline.conf.sample; then rm -f %D/etc/dancer-ircd/kline.conf; fi
etc/dancer-ircd/kline.conf.sample
@exec if [ ! -f %B/kline.conf ]; then cp %B/%f %B/kline.conf; fi
@unexec if cmp -s %D/etc/dancer-ircd/dline.conf %D/etc/dancer-ircd/dline.conf.sample; then rm -f %D/etc/dancer-ircd/dline.conf; fi
etc/dancer-ircd/dline.conf.sample
@exec if [ ! -f %B/dline.conf ]; then cp %B/%f %B/dline.conf; fi
@unexec if cmp -s %D/etc/dancer-ircd/ohelp %D/etc/dancer-ircd/ohelp.sample; then rm -f %D/etc/dancer-ircd/ohelp; fi
etc/dancer-ircd/ohelp.sample
@exec if [ ! -f %B/ohelp ]; then cp %B/%f %B/ohelp; fi
@unexec if cmp -s %D/etc/dancer-ircd/omotd %D/etc/dancer-ircd/omotd.sample; then rm -f %D/etc/dancer-ircd/omotd; fi
etc/dancer-ircd/omotd.sample
@exec if [ ! -f %B/omotd ]; then cp %B/%f %B/omotd; fi
etc/rc.d/dancer-ircd.sh
@exec install -d -o ircd -g ircd -m 700 /var/log/dancer-ircd /var/run/dancer-ircd
@unexec rmdir %D/etc/dancer-ircd /var/log/dancer-ircd /var/run/dancer-ircd 2>/dev/null || echo "If you are permanently removing this port, you should do a ``rm -rf ${PKG_PREFIX}/etc/dancer-ircd /var/log/dancer-ircd /var/run/dancer-ircd`` to remove any configuration files and log files left." | fmt
%%PORTDOCS%%%%DOCSDIR%%/dancer-oper-guide.txt
%%PORTDOCS%%%%DOCSDIR%%/dancer-user-guide.txt
%%PORTDOCS%%@dirrm %%DOCSDIR%%
