# mail/quickml/files/quickml.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: quickml
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# quickml_enable (bool):	Set to NO by default.

. /etc/rc.subr

name="quickml"
rcvar=`set_rcvar`

command="%%PREFIX%%/sbin/${name}"
command_interpreter="%%RUBY%%"
pidfile="/var/run/${name}.pid"
required_files=%%PREFIX%%/etc/quickmlrc

load_rc_config $name

: ${quickml_enable="NO"}

run_rc_command "$1"
