# net/daemonlogger/files/openrc-daemonlogger.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf.local or /etc/rc.conf to enable daemonlogger:
# daemonlogger_flags (string):           Set to "-d -l" by default.

name="daemonlogger"

: ${daemonlogger_flags:="-d -l /var/log"}

command=%%PREFIX%%/bin/daemonlogger
command_args=" > /dev/null 2>&1"

depend() {
	keyword -shutdown
}
