# sysutils/diskprep/files/diskprep.in
#!/bin/sh
# $Aero: util/diskprep/diskprep.rc,v 1.1 2005/08/06 00:50:59 brooks Exp $

# PROVIDE: diskprep
# REQUIRE: disks
# BEFORE: fsck localswap
# KEYWORD: nojail

rc_debug=YES

. /etc/rc.subr

name="diskprep"
rc_var=`set_rcvar`
start_precmd="diskprep_prestart"
stop_cmd=":"
command="/usr/local/sbin/diskprep"

diskprep_prestart()
{
	if [ -n "${diskprep_label}" ]; then
		# If our label exists, we're done.
		if [ -c /dev/label/${curlabel} ]; then
			return 0
		fi
		command_args="$command_args -l $diskprep_label"
	fi
	if [ -n "${diskprep_config}" ]; then
		command_args="$command_args -c $diskprep_config"
	fi
	if [ -z "${diskprep_disk}" -o ! -c "/dev/${diskprep_disk}" ]; then
		err 1 "diskprep enabled, but no disk found!"
	fi
	command_args="$command_args $diskprep_disk"
}

load_rc_config $name

run_rc_command "$1"
