# net/wide-dhcp/pkg-plist
@comment $FreeBSD$
etc/rc.d/wide-dhcps.sh.sample
sbin/dhcpc
sbin/dhcpm
sbin/relay
sbin/dhcps
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README.jis
%%PORTDOCS%%%%DOCSDIR%%/intro.dhcp
%%PORTDOCS%%%%DOCSDIR%%/intro.dhcp.jis
%%PORTDOCS%%@dirrm %%DOCSDIR%%
%%EXAMPLESDIR%%/dhcpdb.pool
%%EXAMPLESDIR%%/dhcpdb.relay
%%EXAMPLESDIR%%/dhcpdb.server
@dirrm %%EXAMPLESDIR%%
@unexec rm -f /var/db/dhcp_cache
