# net/wide-dhcp/pkg-plist
@comment $FreeBSD$
sbin/dhcpc
sbin/dhcpm
sbin/relay
sbin/dhcps
share/doc/dhcp/README
share/doc/dhcp/README.jis
share/doc/dhcp/intro.dhcp
share/doc/dhcp/intro.dhcp.jis
share/examples/dhcp/dhcpdb.pool.sample
share/examples/dhcp/dhcpdb.relay.sample
share/examples/dhcp/dhcpdb.server.sample
etc/rc.d/wide-dhcps.sh.sample
@exec echo "updating /etc/services"; cp /etc/services /etc/services.bak; (grep -v ^dhcps /etc/services.bak ; echo "dhcps		67/udp    # dhcp server") > /etc/services
@unexec echo "updating /etc/services"; cp /etc/services /etc/services.bak; (grep -v ^dhcps /etc/services.bak ) > /etc/services
@exec echo "updating /etc/services"; cp /etc/services /etc/services.bak; (grep -v ^dhcpc /etc/services.bak ; echo "dhcpc		68/udp    # dhcp client") > /etc/services
@unexec echo "updating /etc/services"; cp /etc/services /etc/services.bak; (grep -v ^dhcpc /etc/services.bak ) > /etc/services
@unexec rm -f /var/db/dhcp_cache
@dirrm share/doc/dhcp
@dirrm share/examples/dhcp
