# devel/cvsnt/files/cvslockd.sh.in
#!/bin/sh
#
# $FreeBSD: /tmp/pcvs/ports/devel/cvsnt/files/Attic/cvslockd.sh.in,v 1.3 2010-03-27 00:12:51 dougb Exp $
#
# PROVIDE: cvslockd
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
. /etc/rc.subr

[ -z "$cvslockd_enable" ] && cvslockd_enable="NO"

name=cvslockd
rcvar=`set_rcvar`
command=%%PREFIX%%/bin/${name}

load_rc_config ${name}
run_rc_command "$1"
