# dns/dnscheckengine/pkg-plist
bin/dnscheck-dispatcher
bin/dnscheckengine
etc/rc.d/dnscheck-dispatcher
libexec/dnscheck/dnscheck
libexec/dnscheck/dnscheck-12hourmailer
libexec/dnscheck/dnscheck-dispatcher
libexec/dnscheck/dnscheck-hostsyntax
libexec/dnscheck/dnscheck-preflight
libexec/dnscheck/dnscheck-zonediff
man/man1/dnscheck-zonediff.1.gz
man/man1/dnscheck-dispatcher.1.gz
man/man1/dnscheck-preflight.1.gz
man/man1/dnscheck-hostsyntax.1.gz
man/man1/dnscheck.1.gz
man/man1/dnscheck-12hourmailer.1.gz
%%PERL5_MAN3%%/DNSCheck::Test::Nameserver.3.gz
%%PERL5_MAN3%%/DNSCheck::Config.3.gz
%%PERL5_MAN3%%/DNSCheck.3.gz
%%PERL5_MAN3%%/DNSCheck::Lookup::DNS.3.gz
%%PERL5_MAN3%%/DNSCheck::Test::Consistency.3.gz
%%PERL5_MAN3%%/DNSCheck::Test::Host.3.gz
%%PERL5_MAN3%%/DNSCheck::Test::Connectivity.3.gz
%%PERL5_MAN3%%/DNSCheck::Overview.3.gz
%%PERL5_MAN3%%/DNSCheck::Test::Address.3.gz
%%PERL5_MAN3%%/DNSCheck::Test::DNSSEC.3.gz
%%PERL5_MAN3%%/DNSCheck::Test::SMTP.3.gz
%%PERL5_MAN3%%/DNSCheck::Lookup::ASN.3.gz
%%PERL5_MAN3%%/DNSCheck::Test::Zone.3.gz
%%PERL5_MAN3%%/DNSCheck::Logger.3.gz
%%PERL5_MAN3%%/DNSCheck::Test::Common.3.gz
%%PERL5_MAN3%%/DNSCheck::Test::Delegation.3.gz
%%PERL5_MAN3%%/DNSCheck::Test::SOA.3.gz
%%PERL5_MAN3%%/DNSCheck::Lookup::Resolver.3.gz
%%PERL5_MAN3%%/DNSCheck::Test::Mail.3.gz
%%SITE_PERL%%/DNSCheck.pm
%%SITE_PERL%%/DNSCheck/Test/Address.pm
%%SITE_PERL%%/DNSCheck/Test/Consistency.pm
%%SITE_PERL%%/DNSCheck/Test/Host.pm
%%SITE_PERL%%/DNSCheck/Test/SMTP.pm
%%SITE_PERL%%/DNSCheck/Test/Common.pm
%%SITE_PERL%%/DNSCheck/Test/DNSSEC.pm
%%SITE_PERL%%/DNSCheck/Test/Mail.pm
%%SITE_PERL%%/DNSCheck/Test/SOA.pm
%%SITE_PERL%%/DNSCheck/Test/Connectivity.pm
%%SITE_PERL%%/DNSCheck/Test/Delegation.pm
%%SITE_PERL%%/DNSCheck/Test/Nameserver.pm
%%SITE_PERL%%/DNSCheck/Test/Zone.pm
%%SITE_PERL%%/DNSCheck/Config.pm
%%SITE_PERL%%/DNSCheck/Logger.pm
%%SITE_PERL%%/DNSCheck/Locale.pm
%%SITE_PERL%%/DNSCheck/Overview.pod
%%SITE_PERL%%/DNSCheck/Lookup/ASN.pm
%%SITE_PERL%%/DNSCheck/Lookup/DNS.pm
%%SITE_PERL%%/DNSCheck/Lookup/Resolver.pm
%%SITE_PERL%%/%%PERL_ARCH%%/auto/DNSCheck/.packlist
share/dnscheck/locale/en.yaml
share/dnscheck/locale/sv.yaml
share/dnscheck/config.yaml
share/dnscheck/policy.yaml
