# databases/mdcached/pkg-plist
bin/mdcached
etc/rc.d/mdcached
include/mc_client.h
include/mc_protocol.h
lib/libmdcached.a
lib/libmdcached.so
lib/libmdcached.so.1
man/man1/mdcached.1.gz
