# mail/sentinel/pkg-plist
sbin/sentinel
etc/sentinel.cf.dist
etc/rc.d/sentinel.sh-dist
@dir(smmsp,smmsp,) /var/spool/sentinel
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README.FreeBSD
