# mail/sentinel/pkg-plist
sbin/sentinel
etc/sentinel.cf.dist
etc/rc.d/sentinel.sh-dist
@exec mkdir /var/sentinel
@exec chown smmsp:smmsp /var/sentinel
@unexec rm -rf /var/sentinel
%%PORTDOCS%%share/doc/sentinel/README
%%PORTDOCS%%@dirrm share/doc/sentinel
