# mail/sentinel/pkg-plist
sbin/sentinel
etc/sentinel.cf.dist
etc/rc.d/sentinel.sh-dist
@exec mkdir /var/spool/sentinel
@exec chown smmsp:smmsp /var/spool/sentinel
@unexec rm -rf /var/spool/sentinel
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README.FreeBSD
%%PORTDOCS%%@dirrm %%DOCSDIR%%
