# mail/sentinel/pkg-plist
sbin/sentinel
etc/sentinel.cf.dist
etc/rc.d/000.sentinel.sh
@exec mkdir %D/var/sentinel
@exec chown smmsp.smmsp %D/var/sentinel
@unexec rmdir %D/var/sentinel 2>/dev/null || true
%%PORTDOCS%%share/doc/sentinel/README
%%PORTDOCS%%@dirrm share/doc/sentinel
