# mail/youbin/pkg-plist
bin/youbin
bin/xyoubin
bin/sendto
sbin/youbind
etc/rc.d/youbind.sh.sample
@dirrm sbin
