# dns/pdnsd/files/pdnsd.in
#!/bin/sh

# PROVIDE: pdnsd
# REQUIRE: netif
# BEFORE: LOGIN

. /etc/rc.subr

name="pdnsd"
rcvar=${name}_enable
command="%%PREFIX%%/sbin/pdnsd"
command_args="-d"
required_files="%%PREFIX%%/etc/pdnsd.conf"
load_rc_config $name

run_rc_command "$1"
