# sysutils/webmin/files/webmin.sh.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: webmin
# REQUIRE: LOGIN
#
# Add the following line to /etc/rc.conf to enable webmin:
#
# webmin_enable="YES"
#

webmin_enable=${webmin_enable:-"NO"}

. /etc/rc.subr

name=webmin
rcvar=webmin_enable

procname=%%PERL5%%
pidfile=/var/log/webmin/miniserv.pid
required_dirs=%%PREFIX%%/etc/webmin
command=%%PREFIX%%/etc/webmin/start

load_rc_config ${name}
run_rc_command "$1"
