# www/blogd/pkg-plist
bin/blogd
bin/blogd-add
bin/blogd-delete
bin/blogd-extract
bin/blogd-replace
etc/blogd.css.sample
etc/rc.d/blogd
@unexec if cmp -s %D/etc/blogd.config %D/etc/blogd.config.dist; then rm -f %D/etc/blogd.config; fi
etc/blogd.config.dist
