# net/syncthing/files/openrc-syncthing.in
#!/sbin/openrc-run

: ${syncthing_user:="syncthing"}
: ${syncthing_group:="syncthing"}
: ${syncthing_dir:="/var/db/syncthing"}

command=%%PREFIX%%/bin/syncthing
command_user=${syncthing_user}
pidfile=/var/run/syncthing.pid
name="Syncthing"
command_args="-home=${syncthing_dir} -no-browser ${syncthing_args}"
start_stop_daemon_args="-b -m"

depend() {
	need localmount net
	after bootmisc
}

start_pre() {
	if [ ! -d ${syncthing_dir} ]; then
			install -d -o ${syncthing_user} -g ${syncthing_group} ${syncthing_dir}
	fi
}
