# sysutils/volman/files/openrc-volmand.in
#!/sbin/openrc-run

name="volmand"

pidfile="/var/run/volmand.pid"
command="%%PREFIX%%/sbin/volmand"

depend() {
	keyword -shutdown -jail
}
