# sysutils/volman/files/volmand.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: volmand
# REQUIRE: LOGIN
# KEYWORD: shutdown nojail

. /etc/rc.subr

name="volmand"
rcvar=`set_rcvar`

load_rc_config $name

: ${volmand_enable:="NO"}

command="%%PREFIX%%/sbin/volmand"
pidfile="/var/run/volmand.pid"

run_rc_command "$1"
