# russian/ircd-hybrid/pkg-plist
bin/ircd
bin/viconf
bin/mkpasswd
bin/convertconf
bin/convertilines
bin/convertklines
bin/encspeed
bin/vimotd
bin/viklines
bin/servlink
etc/rc.d/ircd-hybrid-ru.sh.sample
etc/ircd-hybrid-ru/auto.am.conf
etc/ircd-hybrid-ru/auto.az.conf
etc/ircd-hybrid-ru/auto.by.conf
etc/ircd-hybrid-ru/auto.ee.conf
etc/ircd-hybrid-ru/auto.ge.conf
etc/ircd-hybrid-ru/auto.kg.conf
etc/ircd-hybrid-ru/auto.kz.conf
etc/ircd-hybrid-ru/auto.lt.conf
etc/ircd-hybrid-ru/auto.lv.conf
etc/ircd-hybrid-ru/auto.md.conf
etc/ircd-hybrid-ru/auto.ru.conf
etc/ircd-hybrid-ru/auto.tj.conf
etc/ircd-hybrid-ru/auto.tm.conf
etc/ircd-hybrid-ru/auto.ua.conf
etc/ircd-hybrid-ru/auto.uz.conf
etc/ircd-hybrid-ru/example.conf
etc/ircd-hybrid-ru/example.dalnetru.conf
etc/ircd-hybrid-ru/kline.conf
etc/ircd-hybrid-ru/dline.conf
etc/ircd-hybrid-ru/xline.conf
etc/ircd-hybrid-ru/nresv.conf
etc/ircd-hybrid-ru/cresv.conf
etc/ircd-hybrid-ru/codepage/cp866.tab
etc/ircd-hybrid-ru/codepage/iso8859-5.tab
etc/ircd-hybrid-ru/codepage/koi8.tab
etc/ircd-hybrid-ru/codepage/macintosh.tab
etc/ircd-hybrid-ru/codepage/translit.tab
include/ircd-hybrid-7/balloc.h
include/ircd-hybrid-7/channel.h
include/ircd-hybrid-7/channel_mode.h
include/ircd-hybrid-7/class.h
include/ircd-hybrid-7/client.h
include/ircd-hybrid-7/common.h
include/ircd-hybrid-7/config.h
include/ircd-hybrid-7/defaults.h
include/ircd-hybrid-7/event.h
include/ircd-hybrid-7/fdlist.h
include/ircd-hybrid-7/fileio.h
include/ircd-hybrid-7/handlers.h
include/ircd-hybrid-7/hash.h
include/ircd-hybrid-7/hook.h
include/ircd-hybrid-7/hostmask.h
include/ircd-hybrid-7/irc_string.h
include/ircd-hybrid-7/ircd.h
include/ircd-hybrid-7/ircd_defs.h
include/ircd-hybrid-7/ircd_getopt.h
include/ircd-hybrid-7/ircd_handler.h
include/ircd-hybrid-7/ircd_signal.h
include/ircd-hybrid-7/ircdauth.h
include/ircd-hybrid-7/linebuf.h
include/ircd-hybrid-7/list.h
include/ircd-hybrid-7/listener.h
include/ircd-hybrid-7/m_info.h
include/ircd-hybrid-7/m_kline.h
include/ircd-hybrid-7/md5.h
include/ircd-hybrid-7/memory.h
include/ircd-hybrid-7/modules.h
include/ircd-hybrid-7/motd.h
include/ircd-hybrid-7/msg.h
include/ircd-hybrid-7/numeric.h
include/ircd-hybrid-7/packet.h
include/ircd-hybrid-7/parse.h
include/ircd-hybrid-7/patchlevel.h
include/ircd-hybrid-7/res.h
include/ircd-hybrid-7/restart.h
include/ircd-hybrid-7/resv.h
include/ircd-hybrid-7/rsa.h
include/ircd-hybrid-7/s_auth.h
include/ircd-hybrid-7/s_bsd.h
include/ircd-hybrid-7/s_conf.h
include/ircd-hybrid-7/s_debug.h
include/ircd-hybrid-7/s_gline.h
include/ircd-hybrid-7/s_log.h
include/ircd-hybrid-7/s_misc.h
include/ircd-hybrid-7/s_serv.h
include/ircd-hybrid-7/s_stats.h
include/ircd-hybrid-7/s_user.h
include/ircd-hybrid-7/s_zip.h
include/ircd-hybrid-7/scache.h
include/ircd-hybrid-7/send.h
include/ircd-hybrid-7/serno.h
include/ircd-hybrid-7/setup.h
include/ircd-hybrid-7/sprintf_irc.h
include/ircd-hybrid-7/stdinc.h
include/ircd-hybrid-7/supported.h
include/ircd-hybrid-7/tools.h
include/ircd-hybrid-7/translate.h
include/ircd-hybrid-7/vchannel.h
include/ircd-hybrid-7/whowas.h
lib/ircd-hybrid-ru-modules/autoload/m_accept.so
lib/ircd-hybrid-ru-modules/autoload/m_admin.so
lib/ircd-hybrid-ru-modules/autoload/m_alias.so
lib/ircd-hybrid-ru-modules/autoload/m_away.so
lib/ircd-hybrid-ru-modules/autoload/m_capab.so
lib/ircd-hybrid-ru-modules/autoload/m_cburst.so
lib/ircd-hybrid-ru-modules/autoload/m_challenge.so
lib/ircd-hybrid-ru-modules/autoload/m_cjoin.so
lib/ircd-hybrid-ru-modules/autoload/m_clearchan.so
lib/ircd-hybrid-ru-modules/autoload/m_close.so
lib/ircd-hybrid-ru-modules/autoload/m_connect.so
lib/ircd-hybrid-ru-modules/autoload/m_cryptlink.so
lib/ircd-hybrid-ru-modules/autoload/m_dmem.so
lib/ircd-hybrid-ru-modules/autoload/m_drop.so
lib/ircd-hybrid-ru-modules/autoload/m_encap.so
lib/ircd-hybrid-ru-modules/autoload/m_eob.so
lib/ircd-hybrid-ru-modules/autoload/m_force.so
lib/ircd-hybrid-ru-modules/autoload/m_forcenick.so
lib/ircd-hybrid-ru-modules/autoload/m_gline.so
lib/ircd-hybrid-ru-modules/autoload/m_help.so
lib/ircd-hybrid-ru-modules/autoload/m_info.so
lib/ircd-hybrid-ru-modules/autoload/m_invite.so
lib/ircd-hybrid-ru-modules/autoload/m_ison.so
lib/ircd-hybrid-ru-modules/autoload/m_join.so
lib/ircd-hybrid-ru-modules/autoload/m_kline.so
lib/ircd-hybrid-ru-modules/autoload/m_knock.so
lib/ircd-hybrid-ru-modules/autoload/m_links.so
lib/ircd-hybrid-ru-modules/autoload/m_list.so
lib/ircd-hybrid-ru-modules/autoload/m_lljoin.so
lib/ircd-hybrid-ru-modules/autoload/m_llnick.so
lib/ircd-hybrid-ru-modules/autoload/m_locops.so
lib/ircd-hybrid-ru-modules/autoload/m_lusers.so
lib/ircd-hybrid-ru-modules/autoload/m_map.so
lib/ircd-hybrid-ru-modules/autoload/m_motd.so
lib/ircd-hybrid-ru-modules/autoload/m_names.so
lib/ircd-hybrid-ru-modules/autoload/m_nburst.so
lib/ircd-hybrid-ru-modules/autoload/m_oper.so
lib/ircd-hybrid-ru-modules/autoload/m_operwall.so
lib/ircd-hybrid-ru-modules/autoload/m_ojoin.so
lib/ircd-hybrid-ru-modules/autoload/m_pass.so
lib/ircd-hybrid-ru-modules/autoload/m_ping.so
lib/ircd-hybrid-ru-modules/autoload/m_pong.so
lib/ircd-hybrid-ru-modules/autoload/m_post.so
lib/ircd-hybrid-ru-modules/autoload/m_rehash.so
lib/ircd-hybrid-ru-modules/autoload/m_restart.so
lib/ircd-hybrid-ru-modules/autoload/m_resv.so
lib/ircd-hybrid-ru-modules/autoload/m_set.so
lib/ircd-hybrid-ru-modules/autoload/m_svs.so
lib/ircd-hybrid-ru-modules/autoload/m_stats.so
lib/ircd-hybrid-ru-modules/autoload/m_svinfo.so
lib/ircd-hybrid-ru-modules/autoload/m_tburst.so
lib/ircd-hybrid-ru-modules/autoload/m_testline.so
lib/ircd-hybrid-ru-modules/autoload/m_time.so
lib/ircd-hybrid-ru-modules/autoload/m_topic.so
lib/ircd-hybrid-ru-modules/autoload/m_trace.so
lib/ircd-hybrid-ru-modules/autoload/m_unkline.so
lib/ircd-hybrid-ru-modules/autoload/m_user.so
lib/ircd-hybrid-ru-modules/autoload/m_userhost.so
lib/ircd-hybrid-ru-modules/autoload/m_users.so
lib/ircd-hybrid-ru-modules/autoload/m_version.so
lib/ircd-hybrid-ru-modules/autoload/m_wallops.so
lib/ircd-hybrid-ru-modules/autoload/m_who.so
lib/ircd-hybrid-ru-modules/autoload/m_whois.so
lib/ircd-hybrid-ru-modules/autoload/m_whowas.so
lib/ircd-hybrid-ru-modules/autoload/m_xline.so
lib/ircd-hybrid-ru-modules/m_die.so
lib/ircd-hybrid-ru-modules/m_kick.so
lib/ircd-hybrid-ru-modules/m_kill.so
lib/ircd-hybrid-ru-modules/m_message.so
lib/ircd-hybrid-ru-modules/m_mode.so
lib/ircd-hybrid-ru-modules/m_nick.so
lib/ircd-hybrid-ru-modules/m_part.so
lib/ircd-hybrid-ru-modules/m_quit.so
lib/ircd-hybrid-ru-modules/m_server.so
lib/ircd-hybrid-ru-modules/m_sjoin.so
lib/ircd-hybrid-ru-modules/m_squit.so
messages/ayb/LC_MESSAGES/ircd-hybrid.mo
messages/custom/LC_MESSAGES/ircd-hybrid.mo
help/users/accept
help/users/admin
help/users/away
help/users/challenge
help/users/cjoin
help/users/cmode
help/users/credits
help/users/error
help/users/help
help/users/index
help/users/info
help/users/invite
help/users/ison
help/users/join
help/users/kick
help/users/knock
help/users/links
help/users/list
help/users/lusers
help/users/motd
help/users/names
help/users/nick
help/users/notice
help/users/oper
help/users/part
help/users/pass
help/users/ping
help/users/pong
help/users/privmsg
help/users/quit
help/users/stats
help/users/time
help/users/topic
help/users/umode
help/users/user
help/users/userhost
help/users/users
help/users/version
help/users/who
help/users/whois
help/users/whowas
help/opers/accept
help/opers/admin
help/opers/away
help/opers/capab
help/opers/cburst
help/opers/challenge
help/opers/cjoin
help/opers/client
help/opers/close
help/opers/cmode
help/opers/connect
help/opers/credits
help/opers/cryptlink
help/opers/die
help/opers/dline
help/opers/dmem
help/opers/drop
help/opers/eob
help/opers/error
help/opers/gline
help/opers/help
help/opers/index
help/opers/info
help/opers/invite
help/opers/ison
help/opers/join
help/opers/kick
help/opers/kill
help/opers/kline
help/opers/knock
help/opers/knockll
help/opers/links
help/opers/list
help/opers/lljoin
help/opers/llnick
help/opers/locops
help/opers/lusers
help/opers/modlist
help/opers/modload
help/opers/modrestart
help/opers/modunload
help/opers/motd
help/opers/names
help/opers/nburst
help/opers/nick
help/opers/notice
help/opers/oper
help/opers/operwall
help/opers/part
help/opers/pass
help/opers/ping
help/opers/pong
help/opers/post
help/opers/privmsg
help/opers/quit
help/opers/rehash
help/opers/restart
help/opers/resv
help/opers/server
help/opers/set
help/opers/sjoin
help/opers/squit
help/opers/stats
help/opers/svinfo
help/opers/testline
help/opers/time
help/opers/topic
help/opers/trace
help/opers/uhelp
help/opers/umode
help/opers/undline
help/opers/ungline
help/opers/unkline
help/opers/unresv
help/opers/user
help/opers/userhost
help/opers/users
help/opers/version
help/opers/wallops
help/opers/who
help/opers/whois
help/opers/whowas
%%PORTDOCS%%%%DOCSDIR%%/RELNOTES
%%PORTDOCS%%%%DOCSDIR%%/RELNOTES.hybrid-ru
%%PORTDOCS%%%%DOCSDIR%%/CIDR.txt
%%PORTDOCS%%%%DOCSDIR%%/LazyLinks.txt
%%PORTDOCS%%%%DOCSDIR%%/challenge.txt
%%PORTDOCS%%%%DOCSDIR%%/guidelines.txt
%%PORTDOCS%%%%DOCSDIR%%/index.txt
%%PORTDOCS%%%%DOCSDIR%%/kline.txt
%%PORTDOCS%%%%DOCSDIR%%/messages.txt
%%PORTDOCS%%%%DOCSDIR%%/modeg.txt
%%PORTDOCS%%%%DOCSDIR%%/modes.txt
%%PORTDOCS%%%%DOCSDIR%%/operguide.txt
%%PORTDOCS%%%%DOCSDIR%%/opermyth.txt
%%PORTDOCS%%%%DOCSDIR%%/resv.txt
%%PORTDOCS%%%%DOCSDIR%%/serverhide.txt
%%PORTDOCS%%%%DOCSDIR%%/vchans.txt
%%PORTDOCS%%%%DOCSDIR%%/whats-new.txt
%%PORTDOCS%%%%DOCSDIR%%/server-version-info
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm include/ircd-hybrid-7
@dirrm lib/ircd-hybrid-ru-modules/autoload
@dirrm lib/ircd-hybrid-ru-modules
@dirrm help/users
@dirrm help/opers
@dirrm help
@dirrm messages/ayb/LC_MESSAGES
@dirrm messages/custom/LC_MESSAGES
@dirrm messages/ayb
@dirrm messages/custom
@dirrm messages
@exec mkdir -m 755 /var/log/ircd-hybrid-ru
@exec chown ircdru:ircdru /var/log/ircd-hybrid-ru
@unexec rmdir /var/log/ircd-hybrid-ru 2>/dev/null || true
@unexec rmdir %D/etc/ircd-hybrid-ru/codepage 2>/dev/null || true
@unexec rmdir %D/etc/ircd-hybrid-ru 2>/dev/null || true
