# russian/ircd-hybrid/pkg-plist
bin/encspeed
bin/ircd
bin/mkpasswd
bin/servlink
etc/ircd-hybrid-ru/codepages/cp866.tab
etc/ircd-hybrid-ru/codepages/iso8859-5.tab
etc/ircd-hybrid-ru/codepages/koi8r.tab
etc/ircd-hybrid-ru/codepages/macintosh.tab
etc/ircd-hybrid-ru/codepages/translit.tab
etc/ircd-hybrid-ru/cresv.conf
etc/ircd-hybrid-ru/dline.conf
etc/ircd-hybrid-ru/example.conf
etc/ircd-hybrid-ru/ircd.motd
etc/ircd-hybrid-ru/kline.conf
etc/ircd-hybrid-ru/nresv.conf
etc/ircd-hybrid-ru/rkline.conf
etc/ircd-hybrid-ru/rxline.conf
etc/ircd-hybrid-ru/xline.conf
@sample etc/rc.d/ircd-hybrid-ru.sh.sample
help/opers/accept
help/opers/admin
help/opers/away
help/opers/capab
help/opers/cburst
help/opers/challenge
help/opers/client
help/opers/close
help/opers/cmode
help/opers/codepage
help/opers/connect
help/opers/cryptlink
help/opers/die
help/opers/dline
help/opers/drop
help/opers/eob
help/opers/error
help/opers/etrace
help/opers/gline
help/opers/hash
help/opers/help
help/opers/index
help/opers/info
help/opers/invite
help/opers/ison
help/opers/join
help/opers/kick
help/opers/kill
help/opers/kline
help/opers/knock
help/opers/knockll
help/opers/links
help/opers/list
help/opers/lljoin
help/opers/llnick
help/opers/locops
help/opers/lusers
help/opers/map
help/opers/modlist
help/opers/modload
help/opers/modrestart
help/opers/modunload
help/opers/motd
help/opers/names
help/opers/nburst
help/opers/nick
help/opers/notice
help/opers/omotd
help/opers/oper
help/opers/operwall
help/opers/part
help/opers/pass
help/opers/ping
help/opers/pong
help/opers/post
help/opers/privmsg
help/opers/quit
help/opers/rehash
help/opers/restart
help/opers/resv
help/opers/rkline
help/opers/rxline
help/opers/server
help/opers/set
help/opers/sjoin
help/opers/squit
help/opers/stats
help/opers/svinfo
help/opers/testgecos
help/opers/testline
help/opers/testmask
help/opers/time
help/opers/topic
help/opers/trace
help/opers/uhelp
help/opers/umode
help/opers/undline
help/opers/ungline
help/opers/unkline
help/opers/unresv
help/opers/unxline
help/opers/user
help/opers/userhost
help/opers/users
help/opers/version
help/opers/wallops
help/opers/who
help/opers/whois
help/opers/whowas
help/opers/xline
help/users/accept
help/users/admin
help/users/away
help/users/challenge
help/users/cmode
help/users/error
help/users/help
help/users/index
help/users/info
help/users/invite
help/users/ison
help/users/join
help/users/kick
help/users/knock
help/users/links
help/users/list
help/users/lusers
help/users/map
help/users/motd
help/users/names
help/users/nick
help/users/notice
help/users/oper
help/users/part
help/users/pass
help/users/ping
help/users/pong
help/users/privmsg
help/users/quit
help/users/stats
help/users/time
help/users/topic
help/users/umode
help/users/user
help/users/userhost
help/users/users
help/users/version
help/users/who
help/users/whois
help/users/whowas
help/opers/classlist.ru
help/opers/user.ru
help/opers/help.ru
help/opers/llnick.ru
help/opers/drop.ru
help/opers/challenge.ru
help/opers/uncapture
help/opers/omotd.ru
help/opers/ltrace
help/opers/modunload.ru
help/opers/userhost.ru
help/opers/admin.ru
help/opers/sjoin.ru
help/opers/nick.ru
help/opers/wallops.ru
help/opers/rxline.ru
help/opers/gline.ru
help/opers/restart.ru
help/opers/ltrace.ru
help/opers/away.ru
help/opers/testgecos.ru
help/opers/svinfo.ru
help/opers/knock.ru
help/opers/cryptlink.ru
help/opers/ojoin
help/opers/version.ru
help/opers/list.ru
help/opers/ungline.ru
help/opers/lusers.ru
help/opers/testmask.ru
help/opers/mkpasswd.ru
help/opers/topic.ru
help/opers/modload.ru
help/opers/etrace.ru
help/opers/clearchan.ru
help/opers/knockll.ru
help/opers/squit.ru
help/opers/error.ru
help/opers/forcejoin
help/opers/hash.ru
help/opers/server.ru
help/opers/pong.ru
help/opers/capab.ru
help/opers/invite.ru
help/opers/trace.ru
help/opers/client.ru
help/opers/mkpasswd
help/opers/ojoin.ru
help/opers/capture
help/opers/cmode.ru
help/opers/operwall.ru
help/opers/modrestart.ru
help/opers/rkline.ru
help/opers/whowas.ru
help/opers/post.ru
help/opers/users.ru
help/opers/opme
help/opers/privmsg.ru
help/opers/ping.ru
help/opers/uhelp.ru
help/opers/motd.ru
help/opers/kick.ru
help/opers/kill.ru
help/opers/unresv.ru
help/opers/forcepart
help/opers/ctrace
help/opers/ctrace.ru
help/opers/clearchan
help/opers/codepage.ru
help/opers/quit.ru
help/opers/eob.ru
help/opers/classlist
help/opers/xline.ru
help/opers/cburst.ru
help/opers/ison.ru
help/opers/info.ru
help/opers/index.ru
help/opers/stats.ru
help/opers/opme.ru
help/opers/who.ru
help/opers/umode.ru
help/opers/unxline.ru
help/opers/undline.ru
help/opers/notice.ru
help/opers/nburst.ru
help/opers/forcejoin.ru
help/opers/uncapture.ru
help/opers/modlist.ru
help/opers/capture.ru
help/opers/die.ru
help/opers/set.ru
help/opers/whois.ru
help/opers/lljoin.ru
help/opers/time.ru
help/opers/testline.ru
help/opers/pass.ru
help/opers/dline.ru
help/opers/part.ru
help/opers/rehash.ru
help/opers/links.ru
help/opers/locops.ru
help/opers/unkline.ru
help/opers/names.ru
help/opers/resv.ru
help/opers/accept.ru
help/opers/forcepart.ru
help/opers/connect.ru
help/opers/kline.ru
help/opers/map.ru
help/opers/close.ru
help/opers/join.ru
help/opers/oper.ru
help/users/users.ru
help/users/part.ru
help/users/cmode.ru
help/users/time.ru
help/users/pass.ru
help/users/map.ru
help/users/error.ru
help/users/topic.ru
help/users/join.ru
help/users/oper.ru
help/users/pattern.ru
help/users/notice.ru
help/users/kick.ru
help/users/knock.ru
help/users/ping.ru
help/users/motd.ru
help/users/ison.ru
help/users/who.ru
help/users/quit.ru
help/users/accept.ru
help/users/codepage.ru
help/users/info.ru
help/users/admin.ru
help/users/names.ru
help/users/privmsg.ru
help/users/links.ru
help/users/challenge.ru
help/users/pong.ru
help/users/whois.ru
help/users/lusers.ru
help/users/codepage
help/users/version.ru
help/users/user.ru
help/users/help.ru
help/users/invite.ru
help/users/stats.ru
help/users/index.ru
help/users/umode.ru
help/users/away.ru
help/users/nick.ru
help/users/list.ru
help/users/whowas.ru
help/users/userhost.ru
help/users/mode.ru
lib/ircd-hybrid-ru-modules/autoload/m_accept.so
lib/ircd-hybrid-ru-modules/autoload/m_admin.so
lib/ircd-hybrid-ru-modules/autoload/m_alias.so
lib/ircd-hybrid-ru-modules/autoload/m_away.so
lib/ircd-hybrid-ru-modules/autoload/m_cap.so
lib/ircd-hybrid-ru-modules/autoload/m_capab.so
lib/ircd-hybrid-ru-modules/autoload/m_cburst.so
lib/ircd-hybrid-ru-modules/autoload/m_challenge.so
lib/ircd-hybrid-ru-modules/autoload/m_close.so
lib/ircd-hybrid-ru-modules/autoload/m_codepage.so
lib/ircd-hybrid-ru-modules/autoload/m_connect.so
lib/ircd-hybrid-ru-modules/autoload/m_cryptlink.so
lib/ircd-hybrid-ru-modules/autoload/m_drop.so
lib/ircd-hybrid-ru-modules/autoload/m_encap.so
lib/ircd-hybrid-ru-modules/autoload/m_eob.so
lib/ircd-hybrid-ru-modules/autoload/m_etrace.so
lib/ircd-hybrid-ru-modules/autoload/m_force.so
lib/ircd-hybrid-ru-modules/autoload/m_gline.so
lib/ircd-hybrid-ru-modules/autoload/m_hash.so
lib/ircd-hybrid-ru-modules/autoload/m_help.so
lib/ircd-hybrid-ru-modules/autoload/m_info.so
lib/ircd-hybrid-ru-modules/autoload/m_invite.so
lib/ircd-hybrid-ru-modules/autoload/m_ison.so
lib/ircd-hybrid-ru-modules/autoload/m_kline.so
lib/ircd-hybrid-ru-modules/autoload/m_knock.so
lib/ircd-hybrid-ru-modules/autoload/m_links.so
lib/ircd-hybrid-ru-modules/autoload/m_list.so
lib/ircd-hybrid-ru-modules/autoload/m_lljoin.so
lib/ircd-hybrid-ru-modules/autoload/m_llnick.so
lib/ircd-hybrid-ru-modules/autoload/m_locops.so
lib/ircd-hybrid-ru-modules/autoload/m_lusers.so
lib/ircd-hybrid-ru-modules/autoload/m_map.so
lib/ircd-hybrid-ru-modules/autoload/m_motd.so
lib/ircd-hybrid-ru-modules/autoload/m_names.so
lib/ircd-hybrid-ru-modules/autoload/m_nburst.so
lib/ircd-hybrid-ru-modules/autoload/m_ojoin.so
lib/ircd-hybrid-ru-modules/autoload/m_omotd.so
lib/ircd-hybrid-ru-modules/autoload/m_oper.so
lib/ircd-hybrid-ru-modules/autoload/m_operwall.so
lib/ircd-hybrid-ru-modules/autoload/m_pass.so
lib/ircd-hybrid-ru-modules/autoload/m_ping.so
lib/ircd-hybrid-ru-modules/autoload/m_pong.so
lib/ircd-hybrid-ru-modules/autoload/m_post.so
lib/ircd-hybrid-ru-modules/autoload/m_rehash.so
lib/ircd-hybrid-ru-modules/autoload/m_restart.so
lib/ircd-hybrid-ru-modules/autoload/m_resv.so
lib/ircd-hybrid-ru-modules/autoload/m_rkline.so
lib/ircd-hybrid-ru-modules/autoload/m_rxline.so
lib/ircd-hybrid-ru-modules/autoload/m_set.so
lib/ircd-hybrid-ru-modules/autoload/m_stats.so
lib/ircd-hybrid-ru-modules/autoload/m_svinfo.so
lib/ircd-hybrid-ru-modules/autoload/m_svs.so
lib/ircd-hybrid-ru-modules/autoload/m_tburst.so
lib/ircd-hybrid-ru-modules/autoload/m_testline.so
lib/ircd-hybrid-ru-modules/autoload/m_testmask.so
lib/ircd-hybrid-ru-modules/autoload/m_time.so
lib/ircd-hybrid-ru-modules/autoload/m_topic.so
lib/ircd-hybrid-ru-modules/autoload/m_trace.so
lib/ircd-hybrid-ru-modules/autoload/m_user.so
lib/ircd-hybrid-ru-modules/autoload/m_userhost.so
lib/ircd-hybrid-ru-modules/autoload/m_users.so
lib/ircd-hybrid-ru-modules/autoload/m_version.so
lib/ircd-hybrid-ru-modules/autoload/m_wallops.so
lib/ircd-hybrid-ru-modules/autoload/m_who.so
lib/ircd-hybrid-ru-modules/autoload/m_whois.so
lib/ircd-hybrid-ru-modules/autoload/m_whowas.so
lib/ircd-hybrid-ru-modules/autoload/m_xline.so
lib/ircd-hybrid-ru-modules/m_die.so
lib/ircd-hybrid-ru-modules/m_join.so
lib/ircd-hybrid-ru-modules/m_kick.so
lib/ircd-hybrid-ru-modules/m_kill.so
lib/ircd-hybrid-ru-modules/m_message.so
lib/ircd-hybrid-ru-modules/m_mode.so
lib/ircd-hybrid-ru-modules/m_nick.so
lib/ircd-hybrid-ru-modules/m_part.so
lib/ircd-hybrid-ru-modules/m_quit.so
lib/ircd-hybrid-ru-modules/m_server.so
lib/ircd-hybrid-ru-modules/m_sjoin.so
lib/ircd-hybrid-ru-modules/m_squit.so
man/man8/ircd.8.gz
%%DOCSDIR%%/CIDR.txt
%%DOCSDIR%%/LazyLinks.txt
%%DOCSDIR%%/RELNOTES
%%DOCSDIR%%/RELNOTES.RU
%%DOCSDIR%%/challenge.txt
%%DOCSDIR%%/guidelines.txt
%%DOCSDIR%%/index.txt
%%DOCSDIR%%/kline.txt
%%DOCSDIR%%/messages.txt
%%DOCSDIR%%/modeg.txt
%%DOCSDIR%%/modes.txt
%%DOCSDIR%%/resv.txt
%%DOCSDIR%%/server-version-info
%%DOCSDIR%%/serverhide.txt
%%DOCSDIR%%/whats-new.txt
@dirrm share/doc/ircd
@dirrm lib/ircd-hybrid-ru-modules/autoload
@dirrm lib/ircd-hybrid-ru-modules
@dirrm help/users
@dirrm help/opers
@dirrm help
@dirrm etc/ircd-hybrid-ru/codepages
@dirrm etc/ircd-hybrid-ru
@owner %%USR%%
@group %%GRP%%
@exec install -d -m 755 -o %%USR%% -g %%GRP%% /var/log/ircd-hybrid-ru
@dirrmtry /var/log/ircd-hybrid-ru
@owner root
@group wheel
