# net/isc-dhcp41-relay/pkg-plist
@comment $FreeBSD$
@unexec %D/etc/rc.d/isc-dhcrelay.sh forcestop 2>/dev/null || true
@unexec %D/etc/rc.d/isc-dhcrelay forcestop 2>/dev/null || true
sbin/dhcrelay
