# www/sams/files/sams.sh.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: sams
# REQUIRE: LOGIN cleanvar mysql
# KEYWORD: shutdown

# Define these sams_* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#       /etc/rc.conf.d/sams
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#
# sams_enable="YES"
# sams_flags="<set as needed>"

sams_enable=${sams_enable-"NO"}
sams_pidfile=${sams_pidfile-"/var/run/samsdaemon.pid"}

. /etc/rc.subr

name="sams"
rcvar=sams_enable
command="%%PREFIX%%/bin/samsdaemon"

load_rc_config $name

pidfile="${sams_pidfile}"

run_rc_command "$1"
