# net/generic-nqs/pkg-plist
etc/rc.d/nqs.sh
bin/nmapmgr
bin/qacct
bin/qalter
bin/qcat
bin/qdel
bin/qdev
bin/qhold
bin/qlimit
bin/qmgr
bin/qmsg
bin/qpr
bin/qresume
bin/qrls
bin/qstat
bin/qsub
bin/qsuspend
sbin/finish_autoinst
sbin/loaddaemon
sbin/lpserver
sbin/netclient
sbin/netdaemon
sbin/netserver
sbin/nqsdaemon
sbin/nqsmkdirs
sbin/nqsmktrans
sbin/pipeclient
sbin/shlexefai
sbin/OLDloaddaemon
sbin/OLDnetdaemon
sbin/OLDnetserver
sbin/OLDnqsdaemon
sbin/OLDnqsmkdirs
sbin/OLDnqsmktrans
lib/nqs/GNQS.VERSION
lib/nqs/General.conf.def
lib/nqs/qmgr.hlp
lib/nqs/nmap/machines
@unexec rm -rf %D/lib/nqs
@unexec rm -rf XXXX
