# mail/postfix-policyd-weight/files/openrc-policyd-weight.in
#!/sbin/openrc-run

name="policyd_weight"
command="%%PREFIX%%/bin/policyd-weight"
pidfile=/var/run/policyd-weight.pid
command_args=start

stop_pre() {
    %%PREFIX%%/bin/policyd-weight -k stop
}

depend() {
	keyword -shutdown
}
