# net/wizd/files/wizd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: wizd
# REQUIRE: LOGIN cleanvar
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable wizd:
# wizd_enable (bool):		Set it to "YES" to enable wizd.
#				Default is "NO".
# wizd_flags (str):		Flags passed to widz on startup.
#				Default is "".
#

. /etc/rc.subr

name="wizd"
rcvar=wizd_enable

command="%%PREFIX%%/bin/$name"
command_args="-d"
required_files="%%PREFIX%%/etc/$name.conf"

load_rc_config $name

: ${wizd_enable="NO"}

run_rc_command "$1"
