# net-mgmt/net-snmp-devel/pkg-plist
bin/encode_keychange
bin/mib2c
bin/snmpbulkget
bin/snmpbulkwalk
bin/snmpcheck
bin/snmpconf
bin/snmpdelta
bin/snmpdf
bin/snmpget
bin/snmpgetnext
bin/snmpinform
bin/snmpnetstat
bin/snmpset
bin/snmpstatus
bin/snmptable
bin/snmptest
bin/snmptranslate
bin/snmptrap
bin/snmpusm
bin/snmpvacm
bin/snmpwalk
bin/tkmib
etc/rc.d/snmpd.sh
include/ucd-snmp/agent_index.h
include/ucd-snmp/agent_read_config.h
include/ucd-snmp/agent_registry.h
include/ucd-snmp/agent_trap.h
include/ucd-snmp/asn1.h
include/ucd-snmp/auto_nlist.h
include/ucd-snmp/callback.h
include/ucd-snmp/default_store.h
include/ucd-snmp/ds_agent.h
include/ucd-snmp/header_complex.h
include/ucd-snmp/int64.h
include/ucd-snmp/keytools.h
include/ucd-snmp/mib.h
include/ucd-snmp/mib_module_config.h
include/ucd-snmp/mibincl.h
include/ucd-snmp/parse.h
include/ucd-snmp/read_config.h
include/ucd-snmp/scapi.h
include/ucd-snmp/snmp-tc.h
include/ucd-snmp/snmp.h
include/ucd-snmp/snmp_agent.h
include/ucd-snmp/snmp_alarm.h
include/ucd-snmp/snmp_api.h
include/ucd-snmp/snmp_client.h
include/ucd-snmp/snmp_debug.h
include/ucd-snmp/snmp_impl.h
include/ucd-snmp/snmp_logging.h
include/ucd-snmp/snmp_parse_args.h
include/ucd-snmp/snmp_vars.h
include/ucd-snmp/snmpusm.h
include/ucd-snmp/snmpv3.h
include/ucd-snmp/struct.h
include/ucd-snmp/system.h
include/ucd-snmp/tools.h
include/ucd-snmp/transform_oids.h
include/ucd-snmp/ucd-snmp-agent-includes.h
include/ucd-snmp/ucd-snmp-config.h
include/ucd-snmp/ucd-snmp-includes.h
include/ucd-snmp/util_funcs.h
include/ucd-snmp/var_struct.h
include/ucd-snmp/version.h
@dirrm include/ucd-snmp
lib/libsnmp.a
lib/libsnmp.so
lib/libsnmp.so.%%shlib%%
lib/libucdagent.a
lib/libucdagent.so
lib/libucdagent.so.%%shlib%%
lib/libucdmibs.a
lib/libucdmibs.so
lib/libucdmibs.so.%%shlib%%
sbin/snmpd
sbin/snmptrapd
share/examples/ucd-snmp/EXAMPLE.conf
share/examples/ucd-snmp/ucdDemoPublic.conf
@dirrm share/examples/ucd-snmp
share/snmp/mib2c.conf
share/snmp/mib2c.storage.conf
share/snmp/mib2c.vartypes.conf
share/snmp/mibs/AGENTX-MIB.txt
share/snmp/mibs/DISMAN-SCHEDULE-MIB.txt
share/snmp/mibs/DISMAN-SCRIPT-MIB.txt
share/snmp/mibs/EtherLike-MIB.txt
share/snmp/mibs/HCNUM-TC.txt
share/snmp/mibs/HOST-RESOURCES-MIB.txt
share/snmp/mibs/HOST-RESOURCES-TYPES.txt
share/snmp/mibs/IANA-ADDRESS-FAMILY-NUMBERS-MIB.txt
share/snmp/mibs/IANA-LANGUAGE-MIB.txt
share/snmp/mibs/IANAifType-MIB.txt
share/snmp/mibs/IF-INVERTED-STACK-MIB.txt
share/snmp/mibs/IF-MIB.txt
share/snmp/mibs/INET-ADDRESS-MIB.txt
share/snmp/mibs/IP-MIB.txt
share/snmp/mibs/IP-FORWARD-MIB.txt
share/snmp/mibs/IPV6-ICMP-MIB.txt
share/snmp/mibs/IPV6-MIB.txt
share/snmp/mibs/IPV6-TC.txt
share/snmp/mibs/IPV6-TCP-MIB.txt
share/snmp/mibs/IPV6-UDP-MIB.txt
share/snmp/mibs/RFC1155-SMI.txt
share/snmp/mibs/RFC1213-MIB.txt
share/snmp/mibs/RFC-1215.txt
share/snmp/mibs/RMON-MIB.txt
share/snmp/mibs/SMUX-MIB.txt
share/snmp/mibs/SNMP-COMMUNITY-MIB.txt
share/snmp/mibs/SNMP-FRAMEWORK-MIB.txt
share/snmp/mibs/SNMP-MPD-MIB.txt
share/snmp/mibs/SNMP-NOTIFICATION-MIB.txt
share/snmp/mibs/SNMP-PROXY-MIB.txt
share/snmp/mibs/SNMP-TARGET-MIB.txt
share/snmp/mibs/SNMP-USER-BASED-SM-MIB.txt
share/snmp/mibs/SNMP-VIEW-BASED-ACM-MIB.txt
share/snmp/mibs/SNMPv2-CONF.txt
share/snmp/mibs/SNMPv2-MIB.txt
share/snmp/mibs/SNMPv2-SMI.txt
share/snmp/mibs/SNMPv2-TC.txt
share/snmp/mibs/SNMPv2-TM.txt
share/snmp/mibs/TCP-MIB.txt
share/snmp/mibs/UCD-DEMO-MIB.txt
share/snmp/mibs/UCD-DISKIO-MIB.txt
share/snmp/mibs/UCD-DLMOD-MIB.txt
share/snmp/mibs/UCD-IPFWACC-MIB.txt
share/snmp/mibs/UCD-SNMP-MIB.txt
share/snmp/mibs/UDP-MIB.txt
share/snmp/snmpconf/snmp.conf/authopts
share/snmp/snmpconf/snmp.conf/debugging
share/snmp/snmpconf/snmp.conf/mibs
share/snmp/snmpconf/snmp.conf/output
share/snmp/snmpconf/snmpd.conf/acl
share/snmp/snmpconf/snmpd.conf/basic_setup
share/snmp/snmpconf/snmpd.conf/extending
share/snmp/snmpconf/snmpd.conf/monitor
share/snmp/snmpconf/snmpd.conf/operation
share/snmp/snmpconf/snmpd.conf/system
share/snmp/snmpconf/snmpd.conf/trapsinks
share/snmp/snmpconf/snmptrapd.conf/formatting
share/snmp/snmpconf/snmptrapd.conf/traphandle
@unexec rm -f %D/share/snmp/mibs/.index
@dirrm share/snmp/mibs
@dirrm share/snmp/snmpconf/snmp.conf
@dirrm share/snmp/snmpconf/snmpd.conf
@dirrm share/snmp/snmpconf/snmptrapd.conf
@dirrm share/snmp/snmpconf
@dirrm share/snmp
