# net/pxe/files/openrc-pxe.in
#!/sbin/openrc-run

name=pxe
command="%%PREFIX%%/sbin/pxe"
required_files=%%PREFIX%%/etc/pxe.conf
command_args="-c ${required_files}"
pidfile=/var/run/${name}.pid

depend() {
	use net
	keyword -shutdown
}
