# www/abyssws/files/abyssws.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: abyssws
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# abyssws_enable (bool):	Set to NO by default.
#				Set it to YES to enable abyssws.
#

. /etc/rc.subr

name=abyssws
rcvar=abyssws_enable

command="%%PREFIX%%/lib/abyssws/abyssws"
command_args="-d >/dev/null"
pidfile="%%PREFIX%%/lib/abyssws/log/abyssws.pid"
stop_postcmd="abyssws_poststop"

abyssws_poststop()
{
	/bin/rm -f "${pidfile}"
}

load_rc_config $name

: ${abyssws_enable="NO"}

run_rc_command "$1"
