# mail/vm-pop3d/pkg-plist
etc/rc.d/vm-pop3d.sh.sample
sbin/vm-pop3d
%%PORTDOCS%%share/doc/vm-pop3d/AUTHORS
%%PORTDOCS%%share/doc/vm-pop3d/CHANGES
%%PORTDOCS%%share/doc/vm-pop3d/COPYING
%%PORTDOCS%%share/doc/vm-pop3d/FAQ
%%PORTDOCS%%share/doc/vm-pop3d/INSTALL
%%PORTDOCS%%share/doc/vm-pop3d/README
%%PORTDOCS%%@dirrm share/doc/vm-pop3d
