# net/nakenchat/files/openrc-nakenchat.in
#!/sbin/openrc-run

name="nakenchat"
command=/usr/local/bin/nakenchat
command_args="&"
