# mail/py-ppolicy/files/ppolicyd.in
#!/bin/sh

# Start or stop ppolicyd
# $FreeBSD$

# PROVIDE: ppolicyd
# REQUIRE: TWISTD
# KEYWORD: shutdown

# Define these ppolicyd_* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE

. /etc/rc.subr

name="ppolicyd"
rcvar=ppolicyd_enable

command="%%LOCALBASE%%/bin/twistd"
pidfile="/var/run/ppolicyd.pid"
logfile="/var/log/ppolicyd.log"

load_rc_config $name

: ${ppolicyd_enable="NO"}

command_args="--pidfile=${pidfile} --file=%%PREFIX%%/sbin/ppolilcy.tap --python=%%PREFIX%%/sbin/ppolicy.tap --logfile=${logfile} 2>&1 >/dev/null"
procname="%%PYTHON_CMD%%"

run_rc_command "$1"
