# net-mgmt/netams/files/netams.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: netams
# REQUIRE: NETWORKING mysql
# KEYWORD: shutdown

# Define these netams_* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#       /etc/rc.conf.d/netams
#
# netams_enable (bool):         Set to "NO" by default.
#                               Set it to "YES" to enable netams.
# netams_config (str):          Custom configuration file to be passed
#                               to netams.
# netams_pidfile (str):         Path to the pid file (/var/run/netams.pid)
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE

. /etc/rc.subr

name="netams"
rcvar=`set_rcvar`
command="%%PREFIX%%/libexec/netams"

load_rc_config $name

netams_enable=${netams_enable:-"NO"}
netams_config=${netams_config:-"%%PREFIX%%/etc/netams/netams.conf"}
pidfile="${netams_pidfile:-/var/run/netams.pid}"

required_files="$netams_config"

command_args="-lf ${netams_config} > /dev/null"

run_rc_command "$1"
