# audio/umurmur/files/umurmur.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: umurmur
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# umurmur_enable (bool):	Set to NO by default.
#				Set it to YES to enable umurmur.
# umurmur_flags (str):		Set to "" by default.
#				Extra flags passed to start command.
#

. /etc/rc.subr

# Set some defaults
: ${umurmur_enable="NO"}
: ${umurmur_user="umurmur"}
: ${umurmur_flags=""}

name=umurmur
rcvar=${name}_enable

# pidfile
pidfile="/var/run/umurmur/umurmur.pid"

command="%%PREFIX%%/sbin/umurmurd"
command_args="-p ${pidfile}"

load_rc_config $name

run_rc_command "$1"
