# science/libcifpp/pkg-plist
etc/cron.weekly/update-libcifpp-data
%%ETCDIR%%.conf
include/cif++.hpp
include/cif++/atom_type.hpp
include/cif++/category.hpp
include/cif++/compound.hpp
include/cif++/condition.hpp
include/cif++/datablock.hpp
include/cif++/dictionary_parser.hpp
include/cif++/exports.hpp
include/cif++/file.hpp
include/cif++/format.hpp
include/cif++/forward_decl.hpp
include/cif++/gzio.hpp
include/cif++/item.hpp
include/cif++/iterator.hpp
include/cif++/matrix.hpp
include/cif++/model.hpp
include/cif++/parser.hpp
include/cif++/pdb.hpp
include/cif++/pdb/cif2pdb.hpp
include/cif++/pdb/io.hpp
include/cif++/pdb/pdb2cif.hpp
include/cif++/pdb/tls.hpp
include/cif++/point.hpp
include/cif++/row.hpp
include/cif++/symmetry.hpp
include/cif++/text.hpp
include/cif++/utilities.hpp
include/cif++/validate.hpp
lib/cmake/cifpp/cifpp-config-version.cmake
lib/cmake/cifpp/cifpp-config.cmake
lib/cmake/cifpp/cifpp-targets-%%CMAKE_BUILD_TYPE%%.cmake
lib/cmake/cifpp/cifpp-targets.cmake
lib/libcifpp.so
lib/libcifpp.so.7.0
lib/libcifpp.so.7.0.4
%%DATADIR%%/components.cif
%%DATADIR%%/mmcif_ddl.dic
%%DATADIR%%/mmcif_ma.dic
%%DATADIR%%/mmcif_pdbx.dic
