# sysutils/rshim-user-space/pkg-plist
etc/rc.d/rshim
etc/rshim.conf
@comment sbin/bfb-install
sbin/rshim
share/doc/rshim/README.md
share/man/man8/rshim.8.gz
@comment share/man/man8/bfb-install.8.gz
