# sysutils/rshim-user-space/pkg-plist
etc/rc.d/rshim
sbin/rshim
share/doc/rshim/README.md
man/man8/rshim.8.gz
