# dns/dhid/pkg-plist
etc/dhid.conf.sample
etc/rc.d/dhid.sh
sbin/dhid
