# dns/dhid/pkg-plist
bin/dtool
etc/dhid.conf
etc/rc.d/dhid.sh
sbin/dhid
