# audio/rplay/files/openrc-rplayd.in
#!/sbin/openrc-run

name=rplayd
command=%%PREFIX%%/sbin/${name}

depend() {
	keyword -shutdown
}
