# net/asterisk16/files/asterisk.sh.in
#!/bin/sh

# PROVIDE: asterisk
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable asterisk:
#
# asterisk_enable="YES"
#

. /etc/rc.subr

stop_cmd=asterisk_shutdown
asterisk_shutdown () {
  echo 'Stopping asterisk'
  %%PREFIX%%/sbin/asterisk -qrx 'core stop now'
  sleep 1
  return 0
}

name=asterisk
rcvar=`set_rcvar`

command=%%PREFIX%%/sbin/asterisk

load_rc_config $name

pidfile=${asterisk_pidfile:-"/var/run/asterisk.pid"}

asterisk_enable=${asterisk_enable:-"NO"}

run_rc_command "$1"
