# net/asterisk16/files/asterisk.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: asterisk
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable asterisk:
#
# asterisk_enable="YES"
#

. /etc/rc.subr

name=asterisk
rcvar=asterisk_enable

extra_commands="reload"

stop_cmd="asterisk_stop"
reload_cmd="asterisk_reload"

command="%%PREFIX%%/sbin/asterisk"
command_args="-n -U %%ASTERISK_USER%%"
pidfile=${asterisk_pidfile:-"/var/run/asterisk.pid"}

asterisk_stop() {
  echo 'Stopping asterisk'
  $command -nqrx 'core stop now'
}

asterisk_reload() {
  echo 'Reloading asterisk'
  $command -nqrx 'reload'
}

load_rc_config $name

asterisk_enable=${asterisk_enable:-"NO"}

run_rc_command "$1"
