# net/beacon/files/beacon.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: beacon
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable beacon.
#
#	beacon_enable="YES"
#
# Optional settings:
#
#	beacon_program (default: %%PREFIX%%/bin/beacon)
#	beacon_config (default: %%PREFIX%%/etc/beacon.conf)
#	beacon_user (default: daemon)
#	beacon_outputdir (default: %%PREFIX%%/www/data-dist/beacon)
#

. /etc/rc.subr

name="beacon"
rcvar=`set_rcvar`

load_rc_config $name
command=${beacon_program:-%%PREFIX%%/bin/${name}}
config=${beacon_config:-%%PREFIX%%/etc/${name}.conf}
buser=${beacon_user:-nobody}
outputdir="${beacon_outputdir:-%%PREFIX%%/www/data-dist/beacon}"
if [ -f ${config} ]; then
        toutputdir="`awk '$1 == "OUTPUTDIR" {print $3}' ${config}`"
	if [ -d "${toutputdir}" ]; then
		outputdir="${toutputdir}"
	fi
fi
pidfile="${beacon_pidfile:-${outputdir}/multicastbeacon.PID}"
command_interpreter="%%PREFIX%%/bin/perl"

start_cmd="beacon_start"

beacon_start()
{
	echo "Starting ${name}."
	daemon -u ${buser} -f ${command}
}

run_rc_command "$1"
