# chinese/big5fs/pkg-plist
etc/rc.d/big5fs.sh
modules/cd9660.ko
modules/%%MSDOSFSKO%%
modules/ntfs.ko
@unexec rmdir %D/modules 2>/dev/null || true
