# chinese/big5fs/pkg-plist
etc/rc.d/big5fs.sh
modules/big5cd9660.ko
modules/big5msdos.ko
modules/big5ntfs.ko
@unexec rmdir %D/modules 2>/dev/null || true
