# mail/anubis/pkg-plist
@comment $FreeBSD$
sbin/anubis
etc/rc.d/anubis.sh
@unexec install-info --delete %D/info/anubis.info %D/info/dir
info/anubis.info
@exec install-info %D/info/anubis.info %D/info/dir
%%PORTDOCS%%share/doc/anubis/AUTHORS
%%PORTDOCS%%share/doc/anubis/INSTALL
%%PORTDOCS%%share/doc/anubis/NEWS
%%PORTDOCS%%share/doc/anubis/README
%%PORTDOCS%%share/doc/anubis/THANKS
%%PORTDOCS%%share/doc/anubis/TODO
%%PORTDOCS%%@dirrm share/doc/anubis
share/examples/anubis/1anubisrc
share/examples/anubis/2anubisrc
share/examples/anubis/redhat-defaultrc
%%PAM%%share/examples/anubis/pam/anubis
%%PAM%%share/examples/anubis/pam/anubis.allow
%%PAM%%share/examples/anubis/pam/anubis.deny
%%PAM%%@dirrm share/examples/anubis/pam
@dirrm share/examples/anubis
share/locale/fr/LC_MESSAGES/anubis.mo
