# security/isakmpd/files/isakmpd.sh.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ike isakmpd
# REQUIRE: root mountcritlocal
# BEFORE:  DAEMON
# KEYWORD: nojail

. /etc/rc.subr

name="isakmpd"
rcvar=isakmpd_enable
command="%%PREFIX%%/sbin/${name}"
extra_commands="reload"

load_rc_config $name

: ${isakmpd_enable="NO"}

run_rc_command "$1"
