# security/isakmpd/files/openrc-isakmpd.in
#!/sbin/openrc-run

name="isakmpd"
command="%%PREFIX%%/sbin/${name}"

depend() {
		provide ike isakmpd
		need localmount
		keyword -jail
		after bootmisc
}
