# net-im/proxy65/files/proxy65.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: proxy65
# REQUIRE: DAEMON %%JABBER_REQUIRE%%

. /etc/rc.subr

#
# Extra variables:
# proxy65_jid = JID 
# proxy65_config = configuration file

name="proxy65"
rcvar=proxy65_enable
command="%%PREFIX%%/bin/twistd"

load_rc_config $name

proxy65_enable=${proxy65_enable-"NO"}
proxy65_uid=${proxy65_uid-"%%PROXY65_USER%%"}
proxy65_gid=${proxy65_gid-"%%PROXY65_GROUP%%"}
proxy65_logfile=${proxy65_logfile-"/var/log/proxy65.log"}
proxy65_pidfile=${proxy65_pidfile-"/var/run/proxy65.pid"}
proxy65_config=${proxy65_config-"%%PREFIX%%/etc/proxy65.ini"}

required_files="${proxy65_config}"
command_interpreter="%%PYTHON_CMD%%"
pidfile=${proxy65_pidfile}

if [ -n "${proxy65_jid}" ]; then
	proxy65args="--jid ${proxy65_jid}"
fi

command_args="--originalname -u ${proxy65_uid} -g ${proxy65_gid} --pidfile ${pidfile} -l ${proxy65_logfile} proxy65 -c ${proxy65_config} ${proxy65args}"

run_rc_command "$1"
