# news/nntpcache/pkg-plist
etc/nntpcache/nntpcache.access-dist
etc/nntpcache/nntpcache.config-dist
etc/nntpcache/nntpcache.servers-dist
etc/nntpcache/spam.filter
etc/nntpcache/pubring.pgp
etc/nntpcache/VERSION
etc/nntpcache/newshound.conf
@unexec rm -rf etc/nntpcache
etc/rc.d/nntpcached.sh
bin/confused
sbin/nntpcached
sbin/newshound
man/man8/nntpcache.8.gz
man/man8/newshound.8.gz
@unexec rm -rf /var/spool/nntpcache
