# news/nntpcache/pkg-plist
@comment dig /var/spool/nntpcache
@exec mkdir -p /var/spool/nntpcache
bin/confused
bin/mmap_tests
etc/nntpcache/VERSION
@exec cp %D/%F /var/spool/nntpcache
@exec chown -R news.news /var/spool/nntpcache
@exec chmod 700 /var/spool/nntpcache
@exec chmod 600 /var/spool/nntpcache/VERSION
etc/nntpcache/http/moose.gif
etc/nntpcache/http/nntpcache.jpg
etc/nntpcache/http/index.html
etc/nntpcache/innreport.conf-dist
etc/nntpcache/innreport.pl
etc/nntpcache/innreport.pm
etc/nntpcache/newshound.conf-dist
etc/nntpcache/nntpcache.access-dist
etc/nntpcache/nntpcache.config-dist
etc/nntpcache/nntpcache.servers-dist
etc/nntpcache/pubring.pgp-dist
etc/nntpcache/pubring.asc-dist
etc/nntpcache/spam.filter
etc/rc.d/nntpcached.sh.sample
include/mmap_results.h
libexec/nntpcache-gdb.sh
libexec/unixauth
sbin/innreport.sh
sbin/newshound
sbin/nntpcached
share/doc/nntpcache/AUTHORS
share/doc/nntpcache/ChangeLog
share/doc/nntpcache/FAQ
share/doc/nntpcache/FAQ.html
share/doc/nntpcache/HACKING
share/doc/nntpcache/INSTALL
share/doc/nntpcache/LICENSING
share/doc/nntpcache/NEWS
share/doc/nntpcache/README
share/doc/nntpcache/README.INN
share/doc/nntpcache/VERSION
@exec mkdir -p %D/etc/nntpcache/http/innreport/pics
@dirrm etc/nntpcache/http/innreport/pics
@dirrm etc/nntpcache/http/innreport
@dirrm etc/nntpcache/http
@dirrm etc/nntpcache
@dirrm share/doc/nntpcache
@unexec rm -rf /var/spool/nntpcache
