# news/nntpcache/pkg-plist
etc/nntpcache/access
etc/nntpcache/config
etc/nntpcache/servers
etc/nntpcache/spam.filter
etc/rc.d/nntpcached.sh
sbin/nntpcached
man/man8/nntpcache.8.gz
@dirrm %D/etc/nntpcache
@unexec rm -rf /var/spool/nntpcache
