# emulators/py-nova/files/nova-novncproxy.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: nova_novncproxy
# REQUIRE: LOGIN
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable nova_novncproxy:
#
# nova_novncproxy_enable="YES"
#
# nova_novncproxy_enable (bool):
#	Set it to "YES" to enable nova_novncproxy.
# 	Default is "NO".
#
# nova_novncproxy_logdir (str):
#	Set it to chagge log directory
#	Default is "/var/log/nova"
#
# nova_novncproxy_args (str):
#	Set it to change command line arguments.
#	Default is "--log-file ${nova_novncproxy_logdir}/nova-novncproxy.log"
#

. /etc/rc.subr

name=nova_novncproxy
rcvar=nova_novncproxy_enable

PATH=%%PREFIX%%/bin:%%PREFIX%%/sbin:$PATH

pidfile="/var/run/nova-novncproxy.pid"
procname="%%PREFIX%%/bin/python2.7"

start_precmd=nova_precmd

load_rc_config $name

: ${nova_novncproxy_enable:="NO"}
: ${nova_novncproxy_logdir:="/var/log/nova"}
: ${nova_novncproxy_args:="--log-file ${nova_novncproxy_logdir}/nova-novncproxy.log"}

command="/usr/sbin/daemon"
command_args="-f -p ${pidfile} nova-novncproxy ${nova_novncproxy_args}"

nova_precmd() {
	mkdir -p ${nova_novncproxy_logdir}
}

run_rc_command "$1"
