# databases/tokyotyrant/files/openrc-ttserver.in
#!/sbin/openrc-run

name="ttserver"
command="%%PREFIX%%/bin/${name}"
pidfile="/var/run/ttserver.pid"
: ${ttserver_flags="-dmn -pid ${pidfile}"}
command_args=${ttserver_flags}

depend() {
	keyword -shutdown
}
