# devel/pcsc-lite/pkg-plist
bin/formaticc
bin/installifd
%%PORTDOCS%%share/doc/pcsc-lite/README.DAEMON
%%PORTDOCS%%share/doc/pcsc-lite/pcsc-lite-0.8.7.pdf
%%PORTDOCS%%share/doc/pcsc-lite/muscle-api-1.3.0.pdf
etc/rc.d/pcscd.sh.sample
etc/reader.conf.sample
include/debuglog.h
include/mscdefines.h
include/musclecard.h
include/pcsclite.h
include/winscard.h
include/wintypes.h
lib/libpcsclite.a
lib/libpcsclite.so
lib/libpcsclite.so.0
sbin/pcscd
%%PORTDOCS%%@dirrm share/doc/pcsc-lite
