# devel/pcsc-lite/pkg-plist
bin/formaticc
bin/installifd
%%PORTDOCS%%share/doc/pcsc-lite/README.DAEMON
%%PORTDOCS%%share/doc/pcsc-lite/pcsc-lite-0.8.7.pdf
etc/rc.d/pcscd.sh.sample
etc/reader.conf.sample
include/atrhandler.h
include/configfile.h
include/debuglog.h
include/dyn_generic.h
include/eventhandler.h
include/hotplug.h
include/ifdhandler.h
include/ifdwrapper.h
include/pcsclite.h
include/prothandler.h
include/readerfactory.h
include/sys_generic.h
include/thread_generic.h
include/winscard.h
include/winscard_msg.h
include/winscard_svc.h
include/wintypes.h
lib/libpcsclite.a
lib/libpcsclite.so
lib/libpcsclite.so.0
sbin/pcscd
%%PORTDOCS%%@dirrm share/doc/pcsc-lite
