# www/dansguardian/pkg-plist
@exec mkdir -p etc/dansguardian
@exec mkdir -p etc/dansguardian/languages
@exec mkdir -p etc/dansguardian/languages/arspanish
@exec mkdir -p etc/dansguardian/languages/bulgarian
@exec mkdir -p etc/dansguardian/languages/chinesebig5
@exec mkdir -p etc/dansguardian/languages/chineseb2312
@exec mkdir -p etc/dansguardian/languages/czech
@exec mkdir -p etc/dansguardian/languages/danish
@exec mkdir -p etc/dansguardian/languages/french
@exec mkdir -p etc/dansguardian/languages/german
@exec mkdir -p etc/dansguardian/languages/german2
@exec mkdir -p etc/dansguardian/languages/hebrew
@exec mkdir -p etc/dansguardian/languages/indonesian
@exec mkdir -p etc/dansguardian/languages/italian
@exec mkdir -p etc/dansguardian/languages/lithuanian
@exec mkdir -p etc/dansguardian/languages/mxspanish
@exec mkdir -p etc/dansguardian/languages/polish
@exec mkdir -p etc/dansguardian/languages/portuguese
@exec mkdir -p etc/dansguardian/languages/swedish
@exec mkdir -p etc/dansguardian/languages/turkish
@exec mkdir -p etc/dansguardian/languages/ukenglish
@exec mkdir -p etc/dansguardian/samples
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/chat
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/drugadvocacy
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/dutchbadwords
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/frenchbadwords
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/frenchpornography
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/gambling
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/germanpornography
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/goodphrases
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/googlesearches
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/gore
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/illegaldrugs
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/intolerance
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/italianbadwords
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/legaldrugs
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/personals
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/pornography
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/violence
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/warezhacking
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/weapons
etc/dansguardian/samples/bannedextensionlist
etc/dansguardian/samples/bannediplist
etc/dansguardian/samples/bannedmimetypelist
etc/dansguardian/samples/bannedphraselist
etc/dansguardian/samples/bannedregexpurllist
etc/dansguardian/samples/bannedsitelist
etc/dansguardian/samples/bannedurllist
etc/dansguardian/samples/banneduserlist
etc/dansguardian/samples/contentregexplist
etc/dansguardian/samples/dansguardian.conf
etc/dansguardian/samples/dansguardianf1.conf
etc/dansguardian/samples/exceptioniplist
etc/dansguardian/samples/exceptionphraselist
etc/dansguardian/samples/exceptionsitelist
etc/dansguardian/samples/exceptionurllist
etc/dansguardian/samples/exceptionuserlist
etc/dansguardian/samples/filtergroupslist
etc/dansguardian/samples/greysitelist
etc/dansguardian/samples/greyurllist
etc/dansguardian/samples/pics
etc/dansguardian/samples/transparent1x1.gif
etc/dansguardian/samples/weightedphraselist
etc/dansguardian/samples/dansguardian.pl
etc/dansguardian/languages/ReadMe
etc/dansguardian/languages/arspanish/messages
etc/dansguardian/languages/arspanish/template.html
etc/dansguardian/languages/bulgarian/messages
etc/dansguardian/languages/bulgarian/template.html
etc/dansguardian/languages/chinesebig5/messages
etc/dansguardian/languages/chinesebig5/template.html
etc/dansguardian/languages/chinesegb2312/messages
etc/dansguardian/languages/chinesegb2312/template.html
etc/dansguardian/languages/czech/messages
etc/dansguardian/languages/czech/template.html
etc/dansguardian/languages/danish/messages
etc/dansguardian/languages/danish/template.html
etc/dansguardian/languages/dutch/messages
etc/dansguardian/languages/dutch/template.html
etc/dansguardian/languages/french/messages
etc/dansguardian/languages/french/template.html
etc/dansguardian/languages/german/messages
etc/dansguardian/languages/german/template.html
etc/dansguardian/languages/german2/messages
etc/dansguardian/languages/german2/template.html
etc/dansguardian/languages/hebrew/messages
etc/dansguardian/languages/hebrew/template.html
etc/dansguardian/languages/indonesian/messages
etc/dansguardian/languages/indonesian/template.html
etc/dansguardian/languages/italian/messages
etc/dansguardian/languages/italian/template.html
etc/dansguardian/languages/lithuanian/messages
etc/dansguardian/languages/lithuanian/template.html
etc/dansguardian/languages/mxspanish/messages
etc/dansguardian/languages/mxspanish/template.html
etc/dansguardian/languages/polish/messages
etc/dansguardian/languages/polish/template.html
etc/dansguardian/languages/portuguese/messages
etc/dansguardian/languages/portuguese/template.html
etc/dansguardian/languages/swedish/messages
etc/dansguardian/languages/swedish/template.html
etc/dansguardian/languages/turkish/messages
etc/dansguardian/languages/turkish/template.html
etc/dansguardian/languages/ukenglish/messages
etc/dansguardian/languages/ukenglish/template.html
%%PHRASELISTS%%etc/dansguardian/phraselists/chat/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/drugadvocacy/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/dutchbadwords/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/frenchbadwords/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/frenchpornography/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/gambling/banned
%%PHRASELISTS%%etc/dansguardian/phraselists/gambling/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/germanpornography/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/goodphrases/exception
%%PHRASELISTS%%etc/dansguardian/phraselists/goodphrases/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/googlesearches/banned
%%PHRASELISTS%%etc/dansguardian/phraselists/gore/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/illegaldrugs/banned
%%PHRASELISTS%%etc/dansguardian/phraselists/illegaldrugs/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/intolerance/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/italianbadwords/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/legaldrugs/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/personals/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/pornography/banned
%%PHRASELISTS%%etc/dansguardian/phraselists/pornography/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/violence/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/warezhacking/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/weapons/weighted
etc/rc.d/start-dg.sh
sbin/dansguardian
@dirrm etc/dansguardian/languages/arspanish
@dirrm etc/dansguardian/languages/bulgarian
@dirrm etc/dansguardian/languages/chinesebig5
@dirrm etc/dansguardian/languages/chinesegb2312
@dirrm etc/dansguardian/languages/czech
@dirrm etc/dansguardian/languages/danish
@dirrm etc/dansguardian/languages/dutch
@dirrm etc/dansguardian/languages/french
@dirrm etc/dansguardian/languages/german
@dirrm etc/dansguardian/languages/german2
@dirrm etc/dansguardian/languages/hebrew
@dirrm etc/dansguardian/languages/indonesian
@dirrm etc/dansguardian/languages/italian
@dirrm etc/dansguardian/languages/lithuanian
@dirrm etc/dansguardian/languages/mxspanish
@dirrm etc/dansguardian/languages/polish
@dirrm etc/dansguardian/languages/portuguese
@dirrm etc/dansguardian/languages/swedish
@dirrm etc/dansguardian/languages/turkish
@dirrm etc/dansguardian/languages/ukenglish
@dirrm etc/dansguardian/languages
@dirrm etc/dansguardian/samples
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/chat
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/drugadvocacy
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/dutchbadwords
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/frenchbadwords
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/frenchpornography
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/gambling
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/germanpornography
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/goodphrases
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/googlesearches
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/gore
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/illegaldrugs
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/intolerance
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/italianbadwords
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/legaldrugs
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/personals
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/pornography
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/violence
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/warezhacking
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/weapons
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists
@dirrm etc/dansguardian
@unexec echo If permanently deleting this package, /var/log/dansguardian.log  must be removed manually
@unexec echo If permanently deleting this package, custom configuration files must be removed manually
