# www/dansguardian/pkg-plist
@exec mkdir -p etc/dansguardian
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/chat
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/drugadvocacy
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/dutchbadwords
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/frenchbadwords
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/gambling
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/germanpornography
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/goodphrases
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/gore
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/illegaldrugs
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/intolerance
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/italianbadwords
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/legaldrugs
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/personals
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/pornography
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/violence
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/warezhacking
%%PHRASELISTS%%@exec mkdir -p etc/dansguardian/phraselists/weapons
etc/dansguardian/bannedextensionlist.dist
etc/dansguardian/bannediplist.dist
etc/dansguardian/bannedmimetypelist.dist
etc/dansguardian/bannedphraselist.dist
etc/dansguardian/bannedregexpurllist.dist
etc/dansguardian/bannedsitelist.dist
etc/dansguardian/bannedurllist.dist
etc/dansguardian/banneduserlist.dist
etc/dansguardian/contentregexplist.dist
etc/dansguardian/dansguardian.conf.dist
etc/dansguardian/exceptioniplist.dist
etc/dansguardian/exceptionphraselist.dist
etc/dansguardian/exceptionsitelist.dist
etc/dansguardian/exceptionurllist.dist
etc/dansguardian/exceptionuserlist.dist
etc/dansguardian/messages.dist
%%PHRASELISTS%%etc/dansguardian/phraselists/chat/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/drugadvocacy/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/dutchbadwords/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/frenchbadwords/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/gambling/banned
%%PHRASELISTS%%etc/dansguardian/phraselists/gambling/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/germanpornography/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/goodphrases/exception
%%PHRASELISTS%%etc/dansguardian/phraselists/goodphrases/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/gore/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/illegaldrugs/banned
%%PHRASELISTS%%etc/dansguardian/phraselists/illegaldrugs/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/intolerance/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/italianbadwords/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/legaldrugs/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/personals/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/pornography/banned
%%PHRASELISTS%%etc/dansguardian/phraselists/pornography/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/violence/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/warezhacking/weighted
%%PHRASELISTS%%etc/dansguardian/phraselists/weapons/weighted
etc/dansguardian/pics.dist
etc/dansguardian/template.html.dist
etc/dansguardian/weightedphraselist.dist
etc/rc.d/start-dg.sh
sbin/dansguardian
www/cgi-bin/dansguardian.pl.dist
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/chat
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/drugadvocacy
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/dutchbadwords
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/frenchbadwords
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/gambling
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/germanpornography
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/goodphrases
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/gore
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/illegaldrugs
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/intolerance
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/italianbadwords
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/legaldrugs
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/personals
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/pornography
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/violence
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/warezhacking
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists/weapons
%%PHRASELISTS%%@dirrm etc/dansguardian/phraselists
@dirrm etc/dansguardian
@unexec echo If permanently deleting this package, /var/log/dansguardian.log must be removed manually
@unexec echo If permanently deleting this package, custom configuration files must be removed manually
