# www/dansguardian/pkg-plist
@exec mkdir -p etc/dansguardian
etc/dansguardian/dansguardian.conf
etc/dansguardian/bannedphraselist
etc/dansguardian/bannedextensionlist
etc/dansguardian/bannedmimetypelist
etc/dansguardian/bannedurllist
etc/dansguardian/bannedregexpurllist
etc/dansguardian/bannedsitelist
etc/dansguardian/exceptionsitelist
etc/dansguardian/exceptionuserlist
etc/dansguardian/exceptioniplist
etc/dansguardian/pics
etc/rc.d/dansguardian.sh
sbin/dansguardian
www/cgi-bin/dansguardian.pl
@dirrm etc/dansguardian
@unexec echo If permanently deleting this package, /var/log/dansguardian.log must be removed manually

