# x11/gdm2/pkg-plist
bin/gdm
bin/gdm-binary
bin/gdmXnest
bin/gdmXnestchooser
bin/gdmchooser
bin/gdmflexiserver
bin/gdmgreeter
bin/gdmlogin
bin/gdmphotosetup
bin/gdmsetup
bin/gdmthemetester
etc/gdm/Init/Default
etc/gdm/PostLogin/Default.sample
etc/gdm/PostSession/Default
etc/gdm/PreSession/Default
etc/gdm/XKeepsCrashing
etc/gdm/Xsession
@unexec if cmp -s %D/etc/gdm/gdm.conf %D/etc/gdm/factory-gdm.conf; then rm -f %D/etc/gdm/gdm.conf; fi
etc/gdm/factory-gdm.conf
@exec [ -f %B/gdm.conf ] || cp %B/%f %B/gdm.conf
etc/gdm/locale.alias
@unexec if cmp -s %D/etc/gdm/modules/AccessDwellMouseEvents %D/etc/gdm/modules/factory-AccessDwellMouseEvents; then rm -f %D/etc/gdm/modules/AccessDwellMouseEvents; fi
etc/gdm/modules/factory-AccessDwellMouseEvents
@exec [ -f %B/AccessDwellMouseEvents ] || cp %B/%f %B/AccessDwellMouseEvents
@unexec if cmp -s %D/etc/gdm/modules/AccessKeyMouseEvents %D/etc/gdm/modules/factory-AccessKeyMouseEvents; then rm -f %D/etc/gdm/modules/AccessKeyMouseEvents; fi
etc/gdm/modules/factory-AccessKeyMouseEvents
@exec [ -f %B/AccessKeyMouseEvents ] || cp %B/%f %B/AccessKeyMouseEvents
etc/rc.d/gdm.sh.sample
lib/gtk-2.0/modules/libdwellmouselistener.so
lib/gtk-2.0/modules/libkeymouselistener.so
%%PAM_MISC%%libexec/gdmaskpass
libexec/gdmtranslate
sbin/gdm-restart
sbin/gdm-safe-restart
sbin/gdm-stop
sbin/gdmconfig
share/gnome/applications/gdmflexiserver-xnest.desktop
share/gnome/applications/gdmflexiserver.desktop
share/gnome/applications/gdmsetup.desktop
share/gnome/capplets/gdmphotosetup.desktop
share/gnome/gdm/BuiltInSessions/default.desktop
share/gnome/gdm/gdmchooser.glade
share/gnome/gdm/gdmsetup.glade
share/gnome/gdm/themes/circles/GdmGreeterTheme.desktop
share/gnome/gdm/themes/circles/background.svg
share/gnome/gdm/themes/circles/circles.xml
share/gnome/gdm/themes/circles/flower.png
share/gnome/gdm/themes/circles/help.png
share/gnome/gdm/themes/circles/options.png
share/gnome/gdm/themes/circles/screenshot.png
share/gnome/gdm/themes/happygnome-list/GdmGreeterTheme.desktop
share/gnome/gdm/themes/happygnome-list/background.svg
share/gnome/gdm/themes/happygnome-list/disconnect.png
share/gnome/gdm/themes/happygnome-list/gnome-logo.png
share/gnome/gdm/themes/happygnome-list/happygnome.xml
share/gnome/gdm/themes/happygnome-list/options.png
share/gnome/gdm/themes/happygnome-list/screenshot.png
share/gnome/gdm/themes/happygnome-list/session.png
share/gnome/gdm/themes/happygnome-list/system.png
share/gnome/gdm/themes/happygnome/GdmGreeterTheme.desktop
share/gnome/gdm/themes/happygnome/background.svg
share/gnome/gdm/themes/happygnome/disconnect.png
share/gnome/gdm/themes/happygnome/gnome-logo.png
share/gnome/gdm/themes/happygnome/happygnome.xml
share/gnome/gdm/themes/happygnome/options.png
share/gnome/gdm/themes/happygnome/screenshot.png
share/gnome/gdm/themes/happygnome/session.png
share/gnome/gdm/themes/happygnome/system.png
share/gnome/help/gdm/C/gdm.xml
share/gnome/help/gdm/C/legal.xml
share/gnome/help/gdm/de/figures/gdm_window.png
share/gnome/help/gdm/de/gdm.xml
share/gnome/help/gdm/de/legal.xml
share/gnome/help/gdm/es/figures/gdm_window.png
share/gnome/help/gdm/es/gdm.xml
share/gnome/help/gdm/es/legal.xml
share/gnome/help/gdm/fr/figures/gdm_window.png
share/gnome/help/gdm/fr/gdm.xml
share/gnome/help/gdm/fr/legal.xml
share/gnome/help/gdm/it/figures/gdm_window.png
share/gnome/help/gdm/it/gdm.xml
share/gnome/help/gdm/it/legal.xml
share/gnome/help/gdm/ja/figures/gdm_window.png
share/gnome/help/gdm/ja/gdm.xml
share/gnome/help/gdm/ja/legal.xml
share/gnome/help/gdm/ko/figures/gdm_window.png
share/gnome/help/gdm/ko/gdm.xml
share/gnome/help/gdm/ko/legal.xml
share/gnome/help/gdm/sv/figures/gdm_window.png
share/gnome/help/gdm/sv/gdm.xml
share/gnome/help/gdm/sv/legal.xml
share/gnome/help/gdm/zh_CN/figures/gdm_window.png
share/gnome/help/gdm/zh_CN/gdm.xml
share/gnome/help/gdm/zh_CN/legal.xml
share/gnome/help/gdm/zh_HK/figures/gdm_window.png
share/gnome/help/gdm/zh_HK/gdm.xml
share/gnome/help/gdm/zh_HK/legal.xml
share/gnome/help/gdm/zh_TW/figures/gdm_window.png
share/gnome/help/gdm/zh_TW/gdm.xml
share/gnome/help/gdm/zh_TW/legal.xml
share/gnome/omf/gdm/gdm-C.omf
@exec scrollkeeper-install -q %D/share/gnome/omf/gdm/gdm-C.omf 2>/dev/null || /usr/bin/true
share/gnome/omf/gdm/gdm-de.omf
@exec scrollkeeper-install -q %D/share/gnome/omf/gdm/gdm-de.omf 2>/dev/null || /usr/bin/true
share/gnome/omf/gdm/gdm-es.omf
@exec scrollkeeper-install -q %D/share/gnome/omf/gdm/gdm-es.omf 2>/dev/null || /usr/bin/true
share/gnome/omf/gdm/gdm-fr.omf
@exec scrollkeeper-install -q %D/share/gnome/omf/gdm/gdm-fr.omf 2>/dev/null || /usr/bin/true
share/gnome/omf/gdm/gdm-it.omf
@exec scrollkeeper-install -q %D/share/gnome/omf/gdm/gdm-it.omf 2>/dev/null || /usr/bin/true
share/gnome/omf/gdm/gdm-ja.omf
@exec scrollkeeper-install -q %D/share/gnome/omf/gdm/gdm-ja.omf 2>/dev/null || /usr/bin/true
share/gnome/omf/gdm/gdm-ko.omf
@exec scrollkeeper-install -q %D/share/gnome/omf/gdm/gdm-ko.omf 2>/dev/null || /usr/bin/true
share/gnome/omf/gdm/gdm-sv.omf
@exec scrollkeeper-install -q %D/share/gnome/omf/gdm/gdm-sv.omf 2>/dev/null || /usr/bin/true
share/gnome/omf/gdm/gdm-zh_CN.omf
@exec scrollkeeper-install -q %D/share/gnome/omf/gdm/gdm-zh_CN.omf 2>/dev/null || /usr/bin/true
share/gnome/omf/gdm/gdm-zh_HK.omf
@exec scrollkeeper-install -q %D/share/gnome/omf/gdm/gdm-zh_HK.omf 2>/dev/null || /usr/bin/true
share/gnome/omf/gdm/gdm-zh_TW.omf
@exec scrollkeeper-install -q %D/share/gnome/omf/gdm/gdm-zh_TW.omf 2>/dev/null || /usr/bin/true
share/gnome/pixmaps/gdm-foot-logo.png
share/gnome/pixmaps/gdm-setup.png
share/gnome/pixmaps/gdm-xnest.png
share/gnome/pixmaps/gdm.png
share/gnome/pixmaps/login-photo.png
share/gnome/pixmaps/nobody.png
share/gnome/pixmaps/nohost.png
share/gnome/xsessions/gnome.desktop
share/icons/hicolor/16x16/apps/gdm-xnest.png
share/icons/hicolor/32x32/apps/gdm-setup.png
share/icons/hicolor/32x32/apps/gdm-xnest.png
share/locale/af/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/am/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/ar/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/az/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/be/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/bg/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/bn/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/ca/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/cs/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/cy/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/da/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/de/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/el/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/en_CA/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/en_GB/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/es/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/et/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/eu/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/fa/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/fi/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/fr/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/ga/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/gl/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/gu/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/he/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/hi/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/hr/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/hu/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/id/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/is/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/it/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/ja/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/ko/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/lt/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/lv/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/mi/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/mk/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/ml/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/mn/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/mr/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/ms/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/nl/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/nn/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/no/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/pa/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/pl/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/pt/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/pt_BR/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/ro/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/ru/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/sk/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/sl/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/sq/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/sr/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/sr@Latn/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/sv/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/ta/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/th/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/tr/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/uk/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/vi/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/wa/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/zh_CN/LC_MESSAGES/gdm-%%VERSION%%.mo
share/locale/zh_TW/LC_MESSAGES/gdm-%%VERSION%%.mo
@dirrm share/gnome/xsessions
@dirrm share/gnome/omf/gdm
@dirrm share/gnome/help/gdm/zh_TW/figures
@dirrm share/gnome/help/gdm/zh_TW
@dirrm share/gnome/help/gdm/zh_HK/figures
@dirrm share/gnome/help/gdm/zh_HK
@dirrm share/gnome/help/gdm/zh_CN/figures
@dirrm share/gnome/help/gdm/zh_CN
@dirrm share/gnome/help/gdm/sv/figures
@dirrm share/gnome/help/gdm/sv
@dirrm share/gnome/help/gdm/ko/figures
@dirrm share/gnome/help/gdm/ko
@dirrm share/gnome/help/gdm/ja/figures
@dirrm share/gnome/help/gdm/ja
@dirrm share/gnome/help/gdm/it/figures
@dirrm share/gnome/help/gdm/it
@dirrm share/gnome/help/gdm/fr/figures
@dirrm share/gnome/help/gdm/fr
@dirrm share/gnome/help/gdm/es/figures
@dirrm share/gnome/help/gdm/es
@dirrm share/gnome/help/gdm/de/figures
@dirrm share/gnome/help/gdm/de
@dirrm share/gnome/help/gdm/C
@dirrm share/gnome/help/gdm
@dirrm share/gnome/gdm/themes/happygnome
@dirrm share/gnome/gdm/themes/happygnome-list
@dirrm share/gnome/gdm/themes/circles
@dirrm share/gnome/gdm/themes
@dirrm share/gnome/gdm/BuiltInSessions
@dirrm share/gnome/gdm
@dirrm etc/pam.d
@dirrm etc/gdm/modules
@dirrm etc/gdm/Sessions
@dirrm etc/gdm/PreSession
@dirrm etc/gdm/PostSession
@dirrm etc/gdm/PostLogin
@dirrm etc/gdm/Init
@dirrm etc/gdm
@unexec /bin/rmdir %D/etc/dm/Sessions 2>/dev/null || /usr/bin/true
@unexec /bin/rmdir %D/etc/dm 2>/dev/null || /usr/bin/true
@unexec /bin/rmdir %D/share/gnome/log/gdm 2>/dev/null || /usr/bin/true
@unexec /bin/rmdir %D/share/gnome/log 2>/dev/null || /usr/bin/true
@unexec scrollkeeper-uninstall -q %D/share/gnome/omf/gdm/gdm-C.omf 2>/dev/null || /usr/bin/true
@unexec scrollkeeper-uninstall -q %D/share/gnome/omf/gdm/gdm-de.omf 2>/dev/null || /usr/bin/true
@unexec scrollkeeper-uninstall -q %D/share/gnome/omf/gdm/gdm-es.omf 2>/dev/null || /usr/bin/true
@unexec scrollkeeper-uninstall -q %D/share/gnome/omf/gdm/gdm-fr.omf 2>/dev/null || /usr/bin/true
@unexec scrollkeeper-uninstall -q %D/share/gnome/omf/gdm/gdm-it.omf 2>/dev/null || /usr/bin/true
@unexec scrollkeeper-uninstall -q %D/share/gnome/omf/gdm/gdm-ja.omf 2>/dev/null || /usr/bin/true
@unexec scrollkeeper-uninstall -q %D/share/gnome/omf/gdm/gdm-ko.omf 2>/dev/null || /usr/bin/true
@unexec scrollkeeper-uninstall -q %D/share/gnome/omf/gdm/gdm-sv.omf 2>/dev/null || /usr/bin/true
@unexec scrollkeeper-uninstall -q %D/share/gnome/omf/gdm/gdm-zh_CN.omf 2>/dev/null || /usr/bin/true
@unexec scrollkeeper-uninstall -q %D/share/gnome/omf/gdm/gdm-zh_HK.omf 2>/dev/null || /usr/bin/true
@unexec scrollkeeper-uninstall -q %D/share/gnome/omf/gdm/gdm-zh_TW.omf 2>/dev/null || /usr/bin/true
