# net/nfs-ganesha/files/ganesha.in
#!/bin/sh
#
# PROVIDE: ganesha
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr
name="ganesha"
rcvar=ganesha_enable

command="/usr/local/bin/ganesha.nfsd"

start_precmd="ganesha_checkmodules"
restart_precmd="ganesha_checkmodules"

ganesha_checkmodules() {
	list="getfhat fhreadlink fhlink setthreaduid setthreadgid setthreadgroups"
	for module in $list; do
		kldstat | grep -q $module
		if [ $? -ne 0 ]; then
			if [ -e "/boot/modules/${module}.ko" ]; then
				kldload $module
				if [ $? -ne 0 ]; then
					echo "Could not load kernel module ${module}"
					return 1
				fi
			else
				echo "Could not find kernel module ${module}"
				return 1
			fi
		fi
	done
	return 0
}

load_rc_config $name
run_rc_command $1
