# java/jboss3/pkg-plist
bin/jboss3ctl
etc/rc.d/020.jboss3.sh
%%T%%/server/all/deploy/iiop-service.xml
%%T%%/server/all/deploy/jmx-ejb-connector-server.sar
%%T%%/server/all/deploy/jmx-rmi-adaptor.sar
%%T%%/server/all/deploy/jms-service.xml
%%T%%/server/all/deploy/jmx-ejb-adaptor.jar
%%T%%/server/all/deploy/jmx-html-adaptor.sar
%%T%%/server/all/deploy/user-service.xml
%%T%%/server/all/deploy/jbossmq-destinations-service.xml
%%T%%/server/all/deploy/jbossmq-service.xml
%%T%%/server/all/deploy/jboss-xa.rar
%%T%%/server/all/deploy/jms-ra.rar
%%T%%/server/all/deploy/jboss-local-jdbc.rar
%%T%%/server/all/deploy/jca-service.xml
%%T%%/server/all/deploy/hsqldb-service.xml
%%T%%/server/all/deploy/cluster-service.xml
%%T%%/server/all/deploy/jetty-plugin.sar
%%T%%/server/all/deploy/jboss-net.sar
%%T%%/server/all/deploy/counter-service.xml
%%T%%/server/all/deploy/scheduler-service.xml
%%T%%/server/all/deploy/mail-service.xml
%%T%%/server/all/deploy/properties-service.xml
%%T%%/server/all/deploy/ejb-management.jar
%%T%%/server/all/conf/jacorb.properties
%%T%%/server/all/conf/standardjboss.xml
%%T%%/server/all/conf/login-config.xml
%%T%%/server/all/conf/auth.conf
%%T%%/server/all/conf/server.policy
%%T%%/server/all/conf/log4j.xml
%%T%%/server/all/conf/jboss-minimal.xml
%%T%%/server/all/conf/standardjbosscmp-jdbc.xml
%%T%%/server/all/conf/standardjaws.xml
%%T%%/server/all/conf/jndi.properties
%%T%%/server/all/conf/jbossmq-state.xml
%%T%%/server/all/conf/axis-config.xml
%%T%%/server/all/conf/jboss-service.xml
%%T%%/server/all/lib/tt-bytecode.jar
%%T%%/server/all/lib/jbossmx.jar
%%T%%/server/all/lib/jmxtools.jar
%%T%%/server/all/lib/jsse.jar
%%T%%/server/all/lib/counter-plugin.jar
%%T%%/server/all/lib/jbosssx.jar
%%T%%/server/all/lib/jpl-pattern.jar
%%T%%/server/all/lib/activation.jar
%%T%%/server/all/lib/mail.jar
%%T%%/server/all/lib/tyrex.jar
%%T%%/server/all/lib/bcel.jar
%%T%%/server/all/lib/scheduler-plugin-example.jar
%%T%%/server/all/lib/javagroups-2.0.jar
%%T%%/server/all/lib/jboss-j2ee.jar
%%T%%/server/all/lib/axis.jar
%%T%%/server/all/lib/hsqldb-plugin.jar
%%T%%/server/all/lib/jbossmqha.jar
%%T%%/server/all/lib/jpl-util.jar
%%T%%/server/all/lib/jboss-management.jar
%%T%%/server/all/lib/castor.jar
%%T%%/server/all/lib/org.apache.jasper.jar
%%T%%/server/all/lib/jboss-iiop.jar
%%T%%/server/all/lib/log4j.jar
%%T%%/server/all/lib/gnu-regexp.jar
%%T%%/server/all/lib/jnpserver.jar
%%T%%/server/all/lib/jboss-jsr77.jar
%%T%%/server/all/lib/scheduler-plugin.jar
%%T%%/server/all/lib/jacorb.jar
%%T%%/server/all/lib/commons-logging.jar
%%T%%/server/all/lib/xalan.jar
%%T%%/server/all/lib/jnet.jar
%%T%%/server/all/lib/jboss.jar
%%T%%/server/all/lib/jbossha.jar
%%T%%/server/all/lib/hsqldb.jar
%%T%%/server/all/lib/autonumber-plugin.jar
%%T%%/server/all/lib/mail-plugin.jar
%%T%%/server/all/lib/jboss-jca.jar
%%T%%/server/all/lib/jmxri.jar
%%T%%/server/all/lib/properties-plugin.jar
%%T%%/server/all/lib/jts.jar
%%T%%/server/all/lib/wsdl4j.jar
%%T%%/server/all/lib/jaas.jar
%%T%%/server/all/lib/javax.servlet.jar
%%T%%/server/all/lib/jbossmq.jar
%%T%%/server/all/lib/jcert.jar
%%T%%/server/default/lib/jboss-management.jar
%%T%%/server/default/lib/bcel.jar
%%T%%/server/default/lib/jbossmq.jar
%%T%%/server/default/lib/scheduler-plugin.jar
%%T%%/server/default/lib/jaas.jar
%%T%%/server/default/lib/scheduler-plugin-example.jar
%%T%%/server/default/lib/mail.jar
%%T%%/server/default/lib/properties-plugin.jar
%%T%%/server/default/lib/xalan.jar
%%T%%/server/default/lib/mail-plugin.jar
%%T%%/server/default/lib/jpl-util.jar
%%T%%/server/default/lib/tyrex.jar
%%T%%/server/default/lib/jnet.jar
%%T%%/server/default/lib/jcert.jar
%%T%%/server/default/lib/jbossha.jar
%%T%%/server/default/lib/hsqldb-plugin.jar
%%T%%/server/default/lib/jsse.jar
%%T%%/server/default/lib/jboss.jar
%%T%%/server/default/lib/activation.jar
%%T%%/server/default/lib/jboss-jsr77.jar
%%T%%/server/default/lib/jts.jar
%%T%%/server/default/lib/jboss-j2ee.jar
%%T%%/server/default/lib/autonumber-plugin.jar
%%T%%/server/default/lib/jmxri.jar
%%T%%/server/default/lib/jpl-pattern.jar
%%T%%/server/default/lib/jnpserver.jar
%%T%%/server/default/lib/hsqldb.jar
%%T%%/server/default/lib/jbosssx.jar
%%T%%/server/default/lib/log4j.jar
%%T%%/server/default/lib/counter-plugin.jar
%%T%%/server/default/lib/jbossmx.jar
%%T%%/server/default/lib/javax.servlet.jar
%%T%%/server/default/lib/jboss-jca.jar
%%T%%/server/default/lib/org.apache.jasper.jar
%%T%%/server/default/lib/castor.jar
%%T%%/server/default/lib/gnu-regexp.jar
%%T%%/server/default/lib/jmxtools.jar
%%T%%/server/default/deploy/jetty-plugin.sar
%%T%%/server/default/deploy/jmx-ejb-connector-server.sar
%%T%%/server/default/deploy/properties-service.xml
%%T%%/server/default/deploy/jboss-xa.rar
%%T%%/server/default/deploy/jboss-local-jdbc.rar
%%T%%/server/default/deploy/ejb-management.jar
%%T%%/server/default/deploy/mail-service.xml
%%T%%/server/default/deploy/user-service.xml
%%T%%/server/default/deploy/jms-service.xml
%%T%%/server/default/deploy/jms-ra.rar
%%T%%/server/default/deploy/counter-service.xml
%%T%%/server/default/deploy/scheduler-service.xml
%%T%%/server/default/deploy/jbossmq-service.xml
%%T%%/server/default/deploy/jmx-html-adaptor.sar
%%T%%/server/default/deploy/jca-service.xml
%%T%%/server/default/deploy/jmx-rmi-adaptor.sar
%%T%%/server/default/deploy/hsqldb-service.xml
%%T%%/server/default/deploy/jbossmq-destinations-service.xml
%%T%%/server/default/deploy/jmx-ejb-adaptor.jar
%%T%%/server/default/conf/standardjaws.xml
%%T%%/server/default/conf/jboss-minimal.xml
%%T%%/server/default/conf/jbossmq-state.xml
%%T%%/server/default/conf/server.policy
%%T%%/server/default/conf/login-config.xml
%%T%%/server/default/conf/log4j.xml
%%T%%/server/default/conf/standardjbosscmp-jdbc.xml
%%T%%/server/default/conf/jboss-service.xml
%%T%%/server/default/conf/auth.conf
%%T%%/server/default/conf/jndi.properties
%%T%%/server/default/conf/standardjboss.xml
%%T%%/server/minimal/conf/jndi.properties
%%T%%/server/minimal/conf/log4j.xml
%%T%%/server/minimal/conf/jboss-service.xml
%%T%%/server/minimal/lib/jnet.jar
%%T%%/server/minimal/lib/log4j.jar
%%T%%/server/minimal/lib/jnpserver.jar
%%T%%/lib/concurrent.jar
%%T%%/lib/jboss-common.jar
%%T%%/lib/crimson.jar
%%T%%/lib/jaxp.jar
%%T%%/lib/jboss-boot.jar
%%T%%/lib/log4j-boot.jar
%%T%%/lib/jboss-system.jar
%%T%%/lib/getopt.jar
%%T%%/lib/jboss-jmx.jar
%%T%%/docs/dtd/jaws_2_4.dtd
%%T%%/docs/dtd/jboss-web.dtd
%%T%%/docs/dtd/ejb-jar.dtd
%%T%%/docs/dtd/ejb-jar_2_0.dtd
%%T%%/docs/dtd/jaws.dtd
%%T%%/docs/dtd/connector_1_0.dtd
%%T%%/docs/dtd/jbosscmp-jdbc_3_0.dtd
%%T%%/docs/dtd/web-app_2_3.dtd
%%T%%/docs/dtd/web-app_2_2.dtd
%%T%%/docs/dtd/security_config.dtd
%%T%%/docs/dtd/jboss_3_0.dtd
%%T%%/docs/dtd/jboss-web_3_0.dtd
%%T%%/docs/dtd/jboss_2_4.dtd
%%T%%/docs/dtd/application_1_3.dtd
%%T%%/docs/dtd/jaws_3_0.dtd
%%T%%/docs/dtd/jboss.dtd
%%T%%/docs/dtd/application_1_2.dtd
%%T%%/docs/dtd/jboss_xmbean_1_0.dtd
%%T%%/docs/examples/jca/firebird-service.xml
%%T%%/docs/examples/jca/postgres-service.xml
%%T%%/docs/examples/jca/oracle-service.xml
%%T%%/docs/examples/jca/solid-service.xml
%%T%%/docs/examples/jca/hsqldb-service.xml
%%T%%/docs/examples/jca/informix-xa-service.xml
%%T%%/docs/examples/jca/jdatastore-service.xml
%%T%%/docs/examples/jca/mssql-service.xml
%%T%%/docs/examples/jca/db2-service.xml
%%T%%/docs/examples/jca/informix-service.xml
%%T%%/docs/examples/jca/mysql-service.xml
%%T%%/client/jacorb.jar
%%T%%/client/jboss-system-client.jar
%%T%%/client/log4j.jar
%%T%%/client/jboss-j2ee.jar
%%T%%/client/jnp-client.jar
%%T%%/client/jboss-jsr77.jar
%%T%%/client/jmx-rmi-connector-client.jar
%%T%%/client/jboss-client.jar
%%T%%/client/jmx-connector-client-factory.jar
%%T%%/client/jmx-ejb-connector-client.jar
%%T%%/client/gnu-regexp.jar
%%T%%/client/concurrent.jar
%%T%%/client/auth.conf
%%T%%/client/jbosssx-client.jar
%%T%%/client/jaas.jar
%%T%%/client/cacerts
%%T%%/client/jcert.jar
%%T%%/client/jsse.jar
%%T%%/client/jnet.jar
%%T%%/client/jbossmq-client.jar
%%T%%/client/jbossha-client.jar
%%T%%/client/jbossmqha.jar
%%T%%/client/client-config.xml
%%T%%/client/jboss-net-client.jar
%%T%%/client/jboss-iiop-client.jar
%%T%%/client/jboss-common-client.jar
%%T%%/bin/run.jar
%%T%%/bin/shutdown.sh
%%T%%/bin/shutdown.jar
%%T%%/bin/run.sh
%%T%%/bin/jboss_init_redhat.sh
%%T%%/bin/privatekeys
%%T%%/log/stdout.log
%%T%%/log/stderr.log
man/man1/jboss3ctl.1
@dirrm %%T%%/application-deployments
@dirrm %%T%%/server/all/deploy
@dirrm %%T%%/server/all/conf
@dirrm %%T%%/server/all/lib
@dirrm %%T%%/server/all
@dirrm %%T%%/server/default/lib
@dirrm %%T%%/server/default/deploy
@dirrm %%T%%/server/default/conf
@dirrm %%T%%/server/default
@dirrm %%T%%/server/minimal/deploy
@dirrm %%T%%/server/minimal/conf
@dirrm %%T%%/server/minimal/lib
@dirrm %%T%%/server/minimal
@dirrm %%T%%/server
@dirrm %%T%%/lib
@dirrm %%T%%/docs/dtd
@dirrm %%T%%/docs/examples/jca
@dirrm %%T%%/docs/examples
@dirrm %%T%%/docs
@dirrm %%T%%/client
@dirrm %%T%%/bin
@dirrm %%T%%/log
@dirrm %%T%%
