# mail/prayer/files/openrc-prayer.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:

# prayer_chrooted (bool): Set to NO by default.
#	Force prayer to be executed in chrooted
# prayer_prefork (bool): Set to NO by default.
# 	Allow prayer to prefork itself at startup

name="prayer"
command=%%PREFIX%%/sbin/${name}
pidfile="/var/spool/prayer/pid/${name}.pid"

: ${prayer_chrooted="NO"}
: ${prayer_prefork="NO"}

start_pre()
{
	if yesno prayer_chrooted; then	
		echo "Chroot no fully functionnal for now so not activated"
	#	command=/usr/local/sbin/${name}-chroot
	fi
	if ! checkyesno prayer_prefork; then
		command_args="--disable-prefork"
	fi
return 0
}
