# sysutils/uptimed/pkg-plist
bin/uprecords
etc/uptimed.conf-dist
lib/libuptimed.a
lib/libuptimed.so
lib/libuptimed.so.0
sbin/uptimed
etc/rc.d/uptimed.sh.sample
www/cgi-bin/uprecords.cgi
www/cgi-bin/uprecords.conf
www/cgi-bin/uprecords.footer
www/cgi-bin/uprecords.header
@unexec rmdir %D/www/cgi-bin 2>/dev/null || true
