# sysutils/boxbackup/files/openrc-bbstored.in
#!/sbin/openrc-run

name="bbstored"

: ${bbstored_flags:="%%ETCDIR%%/bbstored.conf"}

pidfile=${bbstored_pidfile:-"/var/run/bbstored.pid"}

command="%%PREFIX%%/sbin/bbstored"

depend() {
	keyword -shutdown
}
