# sysutils/boxbackup/files/openrc-bbackupd.in
#!/sbin/openrc-run

name="bbackupd"

: ${bbackupd_flags:="%%PREFIX%%/etc/boxbackup/bbackupd.conf"}

pidfile=${bbackupd_pidfile:-"/var/run/bbackupd.pid"}

command="%%PREFIX%%/sbin/bbackupd"

depend() {
	keyword -shutdown
}
