# net/mldonkey/pkg-plist
%%CORE%%bin/kill_mldonkey
%%CORE%%bin/mlnet
%%CORE%%bin/mlnet-real
%%CORE%%etc/rc.d/mlnet.sh
%%GUICORE%%bin/mlguistarter
%%GUICORE%%bin/mldonkey_previewer
%%GUICORE%%bin/mlnet+gui
%%GUI%%bin/mlchat
%%GUI%%bin/mlgui
%%GUI%%bin/mlim
