# net/rmsg/files/rmsgd.in
#!/bin/sh

# PROVIDE: rmsg
# REQUIRE: DAEMON
# BEFORE:  LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="rmsgd"
rcvar="rmsgd_enable"
command="%%PREFIX%%/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
