# games/uox3/pkg-plist
uox3/uox3
uox3/uox3.wsc
uox3/wsc_changes.txt
uox3/advance.scp
uox3/c-items.cig
uox3/c-npcs.cig
uox3/create.scp
uox3/items.scp
uox3/location.scp
uox3/menus.scp
uox3/misc.scp
uox3/newbie.scp
uox3/npc.scp
uox3/regions.scp
uox3/server.scp
uox3/skills.scp
uox3/speech.scp
uox3/spells.scp
uox3/teleport.scp
uox3/titles.scp
uox3/tracking.scp
uox3/uox3.ico
uox3/uox3.ini
uox3/c-thumbs.cig
uox3/accounts.adm
uox3/readme.txt
uox3/inscribe.gmp
uox3/wtrigrs.scp
uox3/necro.scp
uox3/ntrigrs.scp
uox3/triggers.scp
uox3/envoke.scp
uox3/house.scp
uox3/info.txt
uox3/midilist.txt
uox3/uoxscriptupdate.txt
uox3/docs/copying
uox3/docs/houses.txt
uox3/docs/inscribe.txt
uox3/docs/items.txt
uox3/docs/itemsadd.txt
uox3/docs/mspells.txt
uox3/docs/newbie.txt
uox3/docs/newbiec.txt
uox3/docs/scriptup.txt
uox3/docs/sounds.txt
uox3/docs/tstone.txt
uox3/docs/anthracks28c.txt
uox3/docs/splitreadme.txt
uox3/docs/update.txt
uox3/docs/trigger.html
uox3/docs/midilist.txt
uox3/docs/readme.txt
uox3/docs/layers.txt
uox3/docs/uoxscriptupdate.txt
uox3/docs/scripteamsupdate.txt
uox3/docs/set-types.txt
uox3/docs/commands.html
uox3/docs/settype.html
uox3/docs/set.html
uox3/docs/setnpcai.html
uox3/docs/setowner.html
uox3/docs/setprivpr1.html
uox3/docs/contents.html
uox3/docs/contents1.html
uox3/docs/actions.txt
uox3/docs/readme.html
uox3/docs/triggerdocs.htm
uox3/docs/npc.html
etc/rc.d/uox3.sh
@dirrm uox3/docs
@dirrm uox3
