# sysutils/lmon/pkg-plist
etc/rc.d/lmon.sh.sample
lmon/control.cfg.sample
lmon/control.pl
lmon/lmon.pl
@dirrmtry lmon
