# security/ossec-hids-server/files/ossec-hids.in
#!/bin/sh
# 
# PROVIDE: ossechids
# REQUIRE: DAEMON
# BEFORE:  LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="ossechids"
rcvar=ossechids_enable

load_rc_config $name

: ${ossechids_enable="NO"}
: ${ossechids_user="ossec"}
: ${ossechids_group="ossec"}

start_precmd=%%PRECMD%%
start_cmd="ossechids_command start"
stop_cmd="ossechids_command stop"
restart_cmd="ossechids_command restart"
status_cmd="ossechids_command status"
reload_cmd="ossechids_command reload"

command="%%PREFIX%%/%%PORTNAME%%/bin/ossec-control"
required_files="%%PREFIX%%/%%PORTNAME%%/etc/ossec.conf"
extra_commands="reload"

fts_queue=%%PREFIX%%/%%PORTNAME%%/queue/fts/fts-queue
ig_queue=%%PREFIX%%/%%PORTNAME%%/queue/fts/ig-queue
ossec_log=%%PREFIX%%/%%PORTNAME%%/logs/ossec.log
active_responses_log=%%PREFIX%%/%%PORTNAME%%/logs/active-responses.log

ossechids_start_precmd() {
    # These files are not created by the daemons with the correct
    # ownership, so create them here before starting up the system,
    # if they don't already exist. This is only done for the "local" and
    # "server" installation types.
    if [ ! -e ${fts_queue} ]; then
        touch ${fts_queue}
        chown ${ossechids_user}:${ossechids_group} ${fts_queue}
        chmod 640 ${fts_queue}
    fi
    if [ ! -e ${ig_queue} ]; then
        touch ${ig_queue}
        chown ${ossechids_user}:${ossechids_group} ${ig_queue}
        chmod 640 ${ig_queue}
    fi

    # Ensure logfiles are created with the correct ownership and mode
    for log in ${ossec_log} ${active_responses_log}; do
	if [ ! -e ${log} ]; then
	    touch ${log}
	    chown ${ossechids_user}:${ossechids_group} ${log}
	    chmod 660 ${log}
	fi
    done
}

ossechids_command() {
	${command} ${rc_arg}
}

run_rc_command "$1"
