# games/openttd/files/openrc-openttd.in
#!/sbin/openrc-run
#
# Add the following lines to /etc/rc.conf to enable openttd server:
#
# openttd_flags="" (string):	Set to "" by default
#				Extra flags passed to start command
#
# Note: flags -D and -f are set by default.

name="openttd"
command="%%PREFIX%%/bin/${name}"
command_args="-s null -m null -v null -b null -D -f"

depend() {
	keyword -shutdown
}
