# net-mgmt/nefu/files/nefu.in
#!/bin/sh

# PROVIDE: nefu
# REQUIRE: DAEMON
# BEFORE:  LOGIN
#
# Add the following line to /etc/rc.conf[.local] to enable nefu
#
# nefu_enable (bool):      Set to "NO" by default.
#                          Set it to "YES" to enable nefu.
# nefu_config (str):       Custom config file for nefu.
#                          (default /usr/local/etc/nefu.conf)
# nefu_flags (str):        Custom additional arguments to be passed to nefu
#

. /etc/rc.subr

name="nefu"
rcvar=`set_rcvar`
command=%%PREFIX%%/bin/${name}

load_rc_config $name

# set default
: ${nefu_enable="NO"}
: ${nefu_config="%%PREFIX%%/etc/nefu.conf"}

command_args="-f ${nefu_config}"

run_rc_command "$1"
