# dns/inadyn-troglobit/pkg-plist
etc/rc.d/inadyn
man/man8/inadyn.8.gz
man/man5/inadyn.conf.5.gz
sbin/inadyn
