# security/denyhosts/files/openrc-denyhosts.in
#!/sbin/openrc-run

name="denyhosts"

command="%%PREFIX%%/bin/denyhosts.py"
command_args="--config %%PREFIX%%/etc/denyhosts.conf --daemon"
pidfile="/var/run/${name}.pid"
