# security/crowdsec/files/crowdsec.cron.in
#minute	hour	mday	month	wday	who	command
0       3       *       *       *       root    %%PREFIX%%/libexec/crowdsec/upgrade-hub
