# mail/sympa-elixus/pkg-plist
etc/rc.d/sympa.sh
etc/sympa.conf
etc/wwsympa.conf
sympa/bin/Archive.pm
sympa/bin/Auth.pm
sympa/bin/CAS.pm
sympa/bin/Commands.pm
sympa/bin/Conf.pm
sympa/bin/Language.pm
sympa/bin/Ldap.pm
sympa/bin/List.pm
sympa/bin/Log.pm
sympa/bin/Marc.pm
sympa/bin/Marc/Search.pm
sympa/bin/Message.pm
sympa/bin/SympaTransport.pm
sympa/bin/Version.pm
sympa/bin/X509.pm
sympa/bin/alias_manager.pl
sympa/bin/aliaswrapper
sympa/bin/arc2webarc.pl
sympa/bin/archived.pl
sympa/bin/bounce-lib.pl
sympa/bin/bounced.pl
sympa/bin/bouncequeue
sympa/bin/cookielib.pm
sympa/bin/create_db.Oracle
sympa/bin/create_db.Pg
sympa/bin/create_db.Sybase
sympa/bin/create_db.mysql
sympa/bin/crypt_passwd.pl
sympa/bin/etc/README
sympa/bin/etc/auth.conf
sympa/bin/etc/ca-bundle.crt
sympa/bin/etc/create_list.conf
sympa/bin/etc/create_list_templates/discussion_list/comment
sympa/bin/etc/create_list_templates/discussion_list/config.tpl
sympa/bin/etc/create_list_templates/discussion_list/config.tt2
sympa/bin/etc/create_list_templates/hotline/comment
sympa/bin/etc/create_list_templates/hotline/config.tpl
sympa/bin/etc/create_list_templates/hotline/config.tt2
sympa/bin/etc/create_list_templates/html-news-letter/comment
sympa/bin/etc/create_list_templates/html-news-letter/config.tpl
sympa/bin/etc/create_list_templates/html-news-letter/config.tt2
sympa/bin/etc/create_list_templates/intranet_list/comment
sympa/bin/etc/create_list_templates/intranet_list/config.tpl
sympa/bin/etc/create_list_templates/intranet_list/config.tt2
sympa/bin/etc/create_list_templates/news-letter/comment
sympa/bin/etc/create_list_templates/news-letter/config.tpl
sympa/bin/etc/create_list_templates/news-letter/config.tt2
sympa/bin/etc/create_list_templates/private_working_group/comment
sympa/bin/etc/create_list_templates/private_working_group/config.tpl
sympa/bin/etc/create_list_templates/private_working_group/config.tt2
sympa/bin/etc/create_list_templates/public_web_forum/comment
sympa/bin/etc/create_list_templates/public_web_forum/config.tpl
sympa/bin/etc/create_list_templates/public_web_forum/config.tt2
sympa/bin/etc/edit_list.conf
sympa/bin/etc/global_task_models/chk_cert_expiration.daily.task
sympa/bin/etc/global_task_models/crl_update.daily.task
sympa/bin/etc/global_task_models/eval_bouncers.daily.task
sympa/bin/etc/global_task_models/expire_bounce.daily.task
sympa/bin/etc/global_task_models/process_bouncers.weekly.task
sympa/bin/etc/global_task_models/purge_orphan_bounces.monthly.task
sympa/bin/etc/global_task_models/purge_user_table.monthly.task
sympa/bin/etc/list_aliases.tt2
sympa/bin/etc/list_task_models/expire.yearly.task
sympa/bin/etc/list_task_models/remind.2month.task
sympa/bin/etc/list_task_models/remind.monthly.task
sympa/bin/etc/list_task_models/remind.yearly.task
sympa/bin/etc/list_task_models/sync_include.ttl.task
sympa/bin/etc/mhonarc-ressources
sympa/bin/etc/mime.types
sympa/bin/etc/scenari/access_web_archive.closed
sympa/bin/etc/scenari/access_web_archive.intranet
sympa/bin/etc/scenari/access_web_archive.listmaster
sympa/bin/etc/scenari/access_web_archive.owner
sympa/bin/etc/scenari/access_web_archive.private
sympa/bin/etc/scenari/access_web_archive.public
sympa/bin/etc/scenari/add.auth
sympa/bin/etc/scenari/add.closed
sympa/bin/etc/scenari/add.default
sympa/bin/etc/scenari/add.owner
sympa/bin/etc/scenari/add.owner_notify
sympa/bin/etc/scenari/create_list.intranet
sympa/bin/etc/scenari/create_list.listmaster
sympa/bin/etc/scenari/create_list.public_listmaster
sympa/bin/etc/scenari/d_edit.default
sympa/bin/etc/scenari/d_edit.owner
sympa/bin/etc/scenari/d_edit.private
sympa/bin/etc/scenari/d_edit.public
sympa/bin/etc/scenari/d_read.default
sympa/bin/etc/scenari/d_read.owner
sympa/bin/etc/scenari/d_read.private
sympa/bin/etc/scenari/d_read.public
sympa/bin/etc/scenari/del.auth
sympa/bin/etc/scenari/del.closed
sympa/bin/etc/scenari/del.default
sympa/bin/etc/scenari/del.owner
sympa/bin/etc/scenari/del.owner_notify
sympa/bin/etc/scenari/edit_list.closed
sympa/bin/etc/scenari/edit_list.owner
sympa/bin/etc/scenari/global_remind.listmaster
sympa/bin/etc/scenari/info.default
sympa/bin/etc/scenari/info.open
sympa/bin/etc/scenari/info.private
sympa/bin/etc/scenari/invite.closed
sympa/bin/etc/scenari/invite.default
sympa/bin/etc/scenari/invite.owner
sympa/bin/etc/scenari/invite.private
sympa/bin/etc/scenari/invite.public
sympa/bin/etc/scenari/remind.default
sympa/bin/etc/scenari/remind.listmaster
sympa/bin/etc/scenari/remind.owner
sympa/bin/etc/scenari/review.closed
sympa/bin/etc/scenari/review.default
sympa/bin/etc/scenari/review.intranet
sympa/bin/etc/scenari/review.listmaster
sympa/bin/etc/scenari/review.owner
sympa/bin/etc/scenari/review.private
sympa/bin/etc/scenari/review.public
sympa/bin/etc/scenari/send.closed
sympa/bin/etc/scenari/send.default
sympa/bin/etc/scenari/send.editor
sympa/bin/etc/scenari/send.editorkey
sympa/bin/etc/scenari/send.editorkeyonly
sympa/bin/etc/scenari/send.editorkeyonlyauth
sympa/bin/etc/scenari/send.intranet
sympa/bin/etc/scenari/send.intranetorprivate
sympa/bin/etc/scenari/send.newsletter
sympa/bin/etc/scenari/send.newsletterkeyonly
sympa/bin/etc/scenari/send.private
sympa/bin/etc/scenari/send.private_smime
sympa/bin/etc/scenari/send.privateandeditorkey
sympa/bin/etc/scenari/send.privateandnomultipartoreditorkey
sympa/bin/etc/scenari/send.privatekey
sympa/bin/etc/scenari/send.privatekeyandeditorkeyonly
sympa/bin/etc/scenari/send.privateoreditorkey
sympa/bin/etc/scenari/send.privateorpublickey
sympa/bin/etc/scenari/send.public
sympa/bin/etc/scenari/send.public_nobcc
sympa/bin/etc/scenari/send.publickey
sympa/bin/etc/scenari/send.publicnoattachment
sympa/bin/etc/scenari/send.publicnomultipart
sympa/bin/etc/scenari/subscribe.auth
sympa/bin/etc/scenari/subscribe.auth_notify
sympa/bin/etc/scenari/subscribe.auth_owner
sympa/bin/etc/scenari/subscribe.closed
sympa/bin/etc/scenari/subscribe.default
sympa/bin/etc/scenari/subscribe.intranet
sympa/bin/etc/scenari/subscribe.intranetorowner
sympa/bin/etc/scenari/subscribe.open
sympa/bin/etc/scenari/subscribe.open_notify
sympa/bin/etc/scenari/subscribe.open_quiet
sympa/bin/etc/scenari/subscribe.owner
sympa/bin/etc/scenari/subscribe.smime
sympa/bin/etc/scenari/subscribe.smimeorowner
sympa/bin/etc/scenari/topics_visibility.conceal
sympa/bin/etc/scenari/topics_visibility.default
sympa/bin/etc/scenari/topics_visibility.identified
sympa/bin/etc/scenari/topics_visibility.noconceal
sympa/bin/etc/scenari/unsubscribe.auth
sympa/bin/etc/scenari/unsubscribe.auth_notify
sympa/bin/etc/scenari/unsubscribe.closed
sympa/bin/etc/scenari/unsubscribe.default
sympa/bin/etc/scenari/unsubscribe.open
sympa/bin/etc/scenari/unsubscribe.open_notify
sympa/bin/etc/scenari/unsubscribe.owner
sympa/bin/etc/scenari/visibility.conceal
sympa/bin/etc/scenari/visibility.default
sympa/bin/etc/scenari/visibility.intranet
sympa/bin/etc/scenari/visibility.noconceal
sympa/bin/etc/scenari/visibility.secret
sympa/bin/etc/sympa.wsdl
sympa/bin/etc/topics.conf
sympa/bin/etc/tt2/bye.tt2
sympa/bin/etc/tt2/certif_warning.tt2
sympa/bin/etc/tt2/digest.tt2
sympa/bin/etc/tt2/expire_deletion.tt2
sympa/bin/etc/tt2/expire_warning1.tt2
sympa/bin/etc/tt2/expire_warning2.tt2
sympa/bin/etc/tt2/global_remind.tt2
sympa/bin/etc/tt2/helpfile.tt2
sympa/bin/etc/tt2/info_report.tt2
sympa/bin/etc/tt2/invite.tt2
sympa/bin/etc/tt2/list_created.tt2
sympa/bin/etc/tt2/list_rejected.tt2
sympa/bin/etc/tt2/list_unknown.tt2
sympa/bin/etc/tt2/listmaster_notification.tt2
sympa/bin/etc/tt2/listowner_notification.tt2
sympa/bin/etc/tt2/lists.tt2
sympa/bin/etc/tt2/moderate.tt2
sympa/bin/etc/tt2/modindex.tt2
sympa/bin/etc/tt2/reject.tt2
sympa/bin/etc/tt2/remind.tt2
sympa/bin/etc/tt2/removed.tt2
sympa/bin/etc/tt2/review.tt2
sympa/bin/etc/tt2/sendpasswd.tt2
sympa/bin/etc/tt2/stats_report.tt2
sympa/bin/etc/tt2/subscriber_notification.tt2
sympa/bin/etc/tt2/summary.tt2
sympa/bin/etc/tt2/welcome.tt2
sympa/bin/etc/tt2/x509-user-cert-missing.tt2
sympa/bin/etc/tt2/your_infected_msg.tt2
sympa/bin/etc/wws.tt2/add_request.tt2
sympa/bin/etc/wws.tt2/admin.tt2
sympa/bin/etc/wws.tt2/admin_menu.tt2
sympa/bin/etc/wws.tt2/arc.tt2
sympa/bin/etc/wws.tt2/arc_manage.tt2
sympa/bin/etc/wws.tt2/arc_protect.tt2
sympa/bin/etc/wws.tt2/arcsearch.tt2
sympa/bin/etc/wws.tt2/arcsearch_form.tt2
sympa/bin/etc/wws.tt2/button_footer.tt2
sympa/bin/etc/wws.tt2/button_header.tt2
sympa/bin/etc/wws.tt2/change_email.tt2
sympa/bin/etc/wws.tt2/choosepasswd.tt2
sympa/bin/etc/wws.tt2/close_list.tt2
sympa/bin/etc/wws.tt2/compose_mail.tt2
sympa/bin/etc/wws.tt2/create_list.tt2
sympa/bin/etc/wws.tt2/create_list_request.tt2
sympa/bin/etc/wws.tt2/d_control.tt2
sympa/bin/etc/wws.tt2/d_editfile.tt2
sympa/bin/etc/wws.tt2/d_read.tt2
sympa/bin/etc/wws.tt2/edit_list_request.tt2
sympa/bin/etc/wws.tt2/editfile.tt2
sympa/bin/etc/wws.tt2/editsubscriber.tt2
sympa/bin/etc/wws.tt2/error.tt2
sympa/bin/etc/wws.tt2/get_closed_lists.tt2
sympa/bin/etc/wws.tt2/get_latest_lists.tt2
sympa/bin/etc/wws.tt2/get_pending_lists.tt2
sympa/bin/etc/wws.tt2/help.tt2
sympa/bin/etc/wws.tt2/help_admin.tt2
sympa/bin/etc/wws.tt2/help_editfile.tt2
sympa/bin/etc/wws.tt2/help_editlist.tt2
sympa/bin/etc/wws.tt2/help_user_options.tt2
sympa/bin/etc/wws.tt2/home.tt2
sympa/bin/etc/wws.tt2/info.tt2
sympa/bin/etc/wws.tt2/install_pending_list.tt2
sympa/bin/etc/wws.tt2/list_button_footer.tt2
sympa/bin/etc/wws.tt2/list_button_header.tt2
sympa/bin/etc/wws.tt2/list_menu.tt2
sympa/bin/etc/wws.tt2/lists.tt2
sympa/bin/etc/wws.tt2/login.tt2
sympa/bin/etc/wws.tt2/loginbanner.tt2
sympa/bin/etc/wws.tt2/loginrequest.tt2
sympa/bin/etc/wws.tt2/main.tt2
sympa/bin/etc/wws.tt2/menu.tt2
sympa/bin/etc/wws.tt2/modindex.tt2
sympa/bin/etc/wws.tt2/notice.tt2
sympa/bin/etc/wws.tt2/pref.tt2
sympa/bin/etc/wws.tt2/remindpasswd.tt2
sympa/bin/etc/wws.tt2/remove_arc.tt2
sympa/bin/etc/wws.tt2/rename_list_request.tt2
sympa/bin/etc/wws.tt2/review.tt2
sympa/bin/etc/wws.tt2/reviewbouncing.tt2
sympa/bin/etc/wws.tt2/scenario_test.tt2
sympa/bin/etc/wws.tt2/serveradmin.tt2
sympa/bin/etc/wws.tt2/set_pending_list_request.tt2
sympa/bin/etc/wws.tt2/show_cert.tt2
sympa/bin/etc/wws.tt2/sigrequest.tt2
sympa/bin/etc/wws.tt2/stats.tt2
sympa/bin/etc/wws.tt2/subindex.tt2
sympa/bin/etc/wws.tt2/suboptions.tt2
sympa/bin/etc/wws.tt2/subrequest.tt2
sympa/bin/etc/wws.tt2/title.tt2
sympa/bin/etc/wws.tt2/translate.tt2
sympa/bin/etc/wws.tt2/view_template.tt2
sympa/bin/etc/wws.tt2/view_translations.tt2
sympa/bin/etc/wws.tt2/viewbounce.tt2
sympa/bin/etc/wws.tt2/viewmod.tt2
sympa/bin/etc/wws.tt2/which.tt2
sympa/bin/init_comment.pl
sympa/bin/load_subscribers.pl
sympa/bin/mail.pm
sympa/bin/mod2html.pl
sympa/bin/msg.pl
sympa/bin/p12topem.pl
sympa/bin/parser.pl
sympa/bin/queue
sympa/bin/smtp.pm
sympa/bin/sympa.pl
sympa/bin/sympa_soap_client.pl
sympa/bin/sympa_soap_server.fcgi
sympa/bin/sympa_wizard.pl
sympa/bin/sympasoap.pm
sympa/bin/task_manager.pl
sympa/bin/testldap.pl
sympa/bin/testlogs.pl
sympa/bin/tools.pl
sympa/bin/tt2native.pl
sympa/bin/virtualwrapper
sympa/bin/wwslib.pm
sympa/bin/wwsympa.fcgi
sympa/bin/wwsympa-wrapper
sympa/doc/contents.png
sympa/doc/crossref.png
sympa/doc/img1.old
sympa/doc/img1.png
sympa/doc/img2.old
sympa/doc/img2.png
sympa/doc/index.html
sympa/doc/index.png
sympa/doc/next.png
sympa/doc/next_g.png
sympa/doc/node1.html
sympa/doc/node10.html
sympa/doc/node11.html
sympa/doc/node12.html
sympa/doc/node13.html
sympa/doc/node14.html
sympa/doc/node15.html
sympa/doc/node16.html
sympa/doc/node17.html
sympa/doc/node18.html
sympa/doc/node19.html
sympa/doc/node2.html
sympa/doc/node20.html
sympa/doc/node21.html
sympa/doc/node22.html
sympa/doc/node23.html
sympa/doc/node24.html
sympa/doc/node3.html
sympa/doc/node4.html
sympa/doc/node5.html
sympa/doc/node6.html
sympa/doc/node7.html
sympa/doc/node8.html
sympa/doc/node9.html
sympa/doc/prev.png
sympa/doc/prev_g.png
sympa/doc/sympa.css
sympa/doc/sympa.html
sympa/doc/sympa.pdf
sympa/doc/sympa.ps
sympa/doc/sympa.tex
sympa/doc/up.png
sympa/doc/up_g.png
sympa/nls/cn.msg
sympa/nls/cz.msg
sympa/nls/de.msg
sympa/nls/es.msg
sympa/nls/et.msg
sympa/nls/fi.msg
sympa/nls/fr.msg
sympa/nls/hu.msg
sympa/nls/it.msg
sympa/nls/nl.msg
sympa/nls/pl.msg
sympa/nls/ro.msg
sympa/nls/tw.msg
sympa/nls/us.msg
sympa/po/cn.po
sympa/po/cz.po
sympa/po/de.po
sympa/po/es.po
sympa/po/et.po
sympa/po/fi.po
sympa/po/fr.po
sympa/po/hu.po
sympa/po/it.po
sympa/po/nl.po
sympa/po/pl.po
sympa/po/pt.po
sympa/po/ro.po
sympa/po/tw.po
sympa/po/us.po
sympa/sample/README
sympa/sample/config.digest
sympa/sample/config.include
sympa/sample/config.moderation
sympa/sample/config.privateorpublickey
sympa/sample/config.public
sympa/sample/sample-list/config
sympa/sample/sample-list/info
sympa/sample/sample-list/message.footer
sympa/sample/sample-list/welcome.tpl
sympa/sample/sample-list/welcome.tt2
sympa/sample/sympa.conf
sympa/sample/topics.fr.conf
sympa/sample/topics.hu.conf
sympa/sample/wwsympa.conf
@dirrm sympa/spool/tmp
@dirrm sympa/spool/task
@dirrm sympa/spool/outgoing
@dirrm sympa/spool/msg
@dirrm sympa/spool/moderation
@dirrm sympa/spool/expire
@dirrm sympa/spool/digest
@dirrm sympa/spool/auth
@dirrm sympa/spool
@dirrm sympa/sample/sample-list
@dirrm sympa/sample
@dirrm sympa/po
@dirrm sympa/nls
@dirrm sympa/expl
@dirrm sympa/doc
@dirrm sympa/bin/etc/wws_templates
@dirrm sympa/bin/etc/wws.tt2
@dirrm sympa/bin/etc/tt2
@dirrm sympa/bin/etc/templates
@dirrm sympa/bin/etc/task_models
@dirrm sympa/bin/etc/scenari
@dirrm sympa/bin/etc/list_task_models
@dirrm sympa/bin/etc/global_task_models
@dirrm sympa/bin/etc/general_task_models
@dirrm sympa/bin/etc/create_list_templates/public_web_forum
@dirrm sympa/bin/etc/create_list_templates/private_working_group
@dirrm sympa/bin/etc/create_list_templates/news-letter
@dirrm sympa/bin/etc/create_list_templates/intranet_list
@dirrm sympa/bin/etc/create_list_templates/html-news-letter
@dirrm sympa/bin/etc/create_list_templates/hotline
@dirrm sympa/bin/etc/create_list_templates/discussion_list
@dirrm sympa/bin/etc/create_list_templates
@dirrm sympa/bin/etc
@dirrm sympa/bin/Marc
@dirrm sympa/bin
@dirrm sympa
