# net/chrony/files/openrc-chronyd.in
#!/sbin/openrc-run

name=chronyd
command=%%PREFIX%%/sbin/chronyd
