# sysutils/portupgrade/pkg-plist
etc/rc.d/000.pkgtools.sh
@exec /bin/mkdir -p %D/lib/compat/pkg
@unexec /bin/rmdir %D/lib/compat/pkg 2>/dev/null || true
sbin/pkg_deinstall
sbin/pkg_fetch
sbin/pkg_glob
sbin/pkg_sort
sbin/pkg_which
sbin/pkgdb
sbin/portcvsweb
sbin/portinstall
sbin/ports_glob
sbin/portsclean
sbin/portsdb
sbin/portupgrade
sbin/portversion
%%RUBY_SITELIBDIR%%/pkg.rb
%%RUBY_SITELIBDIR%%/pkgdb.rb
%%RUBY_SITELIBDIR%%/pkginfo.rb
%%RUBY_SITELIBDIR%%/pkgmisc.rb
%%RUBY_SITELIBDIR%%/pkgtools.rb
%%RUBY_SITELIBDIR%%/pkgtsort.rb
%%RUBY_SITELIBDIR%%/pkgversion.rb
%%RUBY_SITELIBDIR%%/portinfo.rb
%%RUBY_SITELIBDIR%%/ports.rb
%%RUBY_SITELIBDIR%%/portsdb.rb
@unexec if cmp -s %D/etc/pkgtools.conf %D/etc/pkgtools.conf.sample; then rm -f %D/etc/pkgtools.conf; fi
etc/pkgtools.conf.sample
@exec [ -f %B/pkgtools.conf ] || cp %B/%f %B/pkgtools.conf
etc/pkgtools.status-pkg.sh
share/zsh/site-functions/_pkgtools
@unexec [ -f %D/bin/zsh ] || rmdir %D/share/zsh/site-functions 2>/dev/null || true
@unexec rmdir %D/share/zsh 2>/dev/null || true
%%PORTDOCS%%share/examples/pkgtools/bash/complete.sample
%%PORTDOCS%%@dirrm share/examples/pkgtools/bash
%%PORTDOCS%%share/examples/pkgtools/tcsh/complete.sample
%%PORTDOCS%%@dirrm share/examples/pkgtools/tcsh
%%PORTDOCS%%@dirrm share/examples/pkgtools
%%PORTDOCS%%share/doc/pkgtools/ChangeLog
%%PORTDOCS%%share/doc/pkgtools/README
%%PORTDOCS%%@dirrm share/doc/pkgtools
