# comms/smstools/pkg-plist
bin/email2sms
bin/sendsms
bin/smsresend
bin/smsd
bin/getsms
bin/putsms
bin/mysmsd
bin/smsevent
etc/rc.d/smsd.sh.sample
etc/smsd.conf
%%EXAMPLESDIR%%/operator_logo2.sms
%%EXAMPLESDIR%%/received_report.sms
%%EXAMPLESDIR%%/received_sms.sms
%%EXAMPLESDIR%%/send_sms.sms
%%EXAMPLESDIR%%/smsd.black
%%EXAMPLESDIR%%/smsd.conf.easy
%%EXAMPLESDIR%%/smsd.conf.full
%%EXAMPLESDIR%%/sql_demo
%%EXAMPLESDIR%%/smsd.conf.easy.orig
%%EXAMPLESDIR%%/smsd.conf.full.orig
%%EXAMPLESDIR%%/smsd.conf.full.bak
%%EXAMPLESDIR%%/smsd.conf.easy.bak
%%EXAMPLESDIR%%/operator_logo1.sms

%%DOCSDIR%%/html/slideshow/status.gif
%%DOCSDIR%%/html/slideshow/eventhandler.gif
%%DOCSDIR%%/html/slideshow/logfile.gif
%%DOCSDIR%%/html/slideshow/modem.gif
%%DOCSDIR%%/html/slideshow/move_l.gif
%%DOCSDIR%%/html/slideshow/move_r.gif
%%DOCSDIR%%/html/slideshow/page1.html
%%DOCSDIR%%/html/slideshow/page2.html
%%DOCSDIR%%/html/slideshow/page3.html
%%DOCSDIR%%/html/slideshow/page4.html
%%DOCSDIR%%/html/slideshow/page5.html
%%DOCSDIR%%/html/slideshow/page6.html
%%DOCSDIR%%/html/slideshow/page7.html
%%DOCSDIR%%/html/slideshow/page8.html
%%DOCSDIR%%/html/slideshow/queue.gif
%%DOCSDIR%%/html/slideshow/sms.gif
%%DOCSDIR%%/html/slideshow/statistic.gif
%%DOCSDIR%%/html/slideshow/blacklist.gif
%%DOCSDIR%%/html/blacklist.html
%%DOCSDIR%%/html/book.html
%%DOCSDIR%%/html/butterfly.gif
%%DOCSDIR%%/html/compiling.html
%%DOCSDIR%%/html/configure.html
%%DOCSDIR%%/html/configure2.html
%%DOCSDIR%%/html/contact.html
%%DOCSDIR%%/html/costs.html
%%DOCSDIR%%/html/download.html
%%DOCSDIR%%/html/email.png
%%DOCSDIR%%/html/book.html~
%%DOCSDIR%%/html/eventhandler.html
%%DOCSDIR%%/html/faq.html
%%DOCSDIR%%/html/fileformat.html
%%DOCSDIR%%/html/gpl.html
%%DOCSDIR%%/html/gpl_de.html
%%DOCSDIR%%/html/gsmmodem.html
%%DOCSDIR%%/html/hardwarecomp.html
%%DOCSDIR%%/html/history.html
%%DOCSDIR%%/html/index.html
%%DOCSDIR%%/html/index.html~
%%DOCSDIR%%/html/key.png
%%DOCSDIR%%/html/license.html
%%DOCSDIR%%/html/links.html
%%DOCSDIR%%/html/logging.html
%%DOCSDIR%%/html/mnp.html
%%DOCSDIR%%/html/news.html
%%DOCSDIR%%/html/public_key_stefan_frings.txt
%%DOCSDIR%%/html/quotes.html
%%DOCSDIR%%/html/redundancy.html
%%DOCSDIR%%/html/referenzen.html
%%DOCSDIR%%/html/run.html
%%DOCSDIR%%/html/softwarecomp.html
%%DOCSDIR%%/html/sql.html
%%DOCSDIR%%/html/sql2.html
%%DOCSDIR%%/html/statformat.html
%%DOCSDIR%%/html/statusmonitor.html
%%DOCSDIR%%/html/tc35.gif
%%DOCSDIR%%/html/thanks.html
%%DOCSDIR%%/html/twist.gif
%%DOCSDIR%%/html/windows.html
%%DOCSDIR%%/html/thanks.gif
%%DOCSDIR%%/html/alarmhandler.html
%%DOCSDIR%%/html/shippingdetails.html
%%DOCSDIR%%/license.txt
%%DOCSDIR%%/manual.html

@dirrm %%EXAMPLESDIR%%
@dirrm %%DOCSDIR%%/html/slideshow
@dirrm %%DOCSDIR%%/html
@dirrm %%DOCSDIR%%
