# irc/thales/files/openrc-thales.in
#!/sbin/openrc-run

name="thales"
command="%%PREFIX%%/sbin/thales"

required_files="%%PREFIX%%/etc/thales.conf"

depend() {
	use net
	keyword -shutdown
}
