# irc/thales/files/thales.sh.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: thales
# REQUIRE: NETWORKING SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable Thales:
# thales_enable (bool):	Set to "NO" by default.
#			Set it to "YES" to enable Thales
#

. /etc/rc.subr

name="thales"
rcvar=thales_enable

command="%%PREFIX%%/sbin/thales"

[ -z "$thales_enable" ]	&& thales_enable="NO"

required_files="%%PREFIX%%/etc/thales.conf"

load_rc_config $name
run_rc_command "$1"
