# multimedia/zoneminder/files/zm.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: zoneminder
# REQUIRE: mysql apache
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# zoneminder_enable (bool):	Set to NO by default.
#				Set it to YES to enable zoneminder.

. /etc/rc.subr

name="zoneminder"
rcvar=${name}_enable

load_rc_config "$name"
: ${zoneminder_enable="NO"}

command=%%PREFIX%%/bin/zmpkg.pl
command_args="$1"
pidfile=/var/run/zm/zm.pid

stop_cmd="zm_stop"
status_cmd="zm_status"

zm_stop() {
  ${command} ${command_args}
}

zm_status() {
  %%PREFIX%%/bin/zmdc.pl status
}

[ "${command_args}" = "faststart" ] && echo -n "zoneminder: 10 seconds pause..." && sleep 10 && echo " done"

run_rc_command "$1"
