# mail/qmail/pkg-plist
@comment $FreeBSD$
%%RCDLINK%%@unexec rm -f %%LOCALBASE%%/etc/rc.d/qmail.sh
%%LDAP%%bin/auth_imap
%%LDAP%%bin/auth_pop
%%LDAP%%bin/auth_smtp
bin/bouncesaying
bin/condredirect
%%LDAP%%bin/condwrite
bin/datemail
%%LDAP%%bin/digest
bin/elq
bin/except
bin/forward
bin/maildir2mbox
bin/maildirmake
bin/maildirwatch
bin/mailsubj
%%LDAP%%bin/pbsadd
%%LDAP%%bin/pbscheck
%%LDAP%%bin/pbsdbd
bin/pinq
bin/predate
bin/preline
bin/qail
bin/qbiff
%%SPAMCONTROL%%bin/qmail-badloadertypes
%%SPAMCONTROL%%bin/qmail-badmimetypes
%%LDAP%%bin/qmail-cdb
bin/qmail-clean
%%LDAP%%bin/qmail-forward
bin/qmail-getpw
%%LDAP%%bin/qmail-group
bin/qmail-inject
%%LDAP%%bin/qmail-ldaplookup
bin/qmail-local
bin/qmail-lspawn
bin/qmail-newmrh
bin/qmail-newu
bin/qmail-pop3d
bin/qmail-popup
bin/qmail-pw2u
bin/qmail-qmqpc
bin/qmail-qmqpd
bin/qmail-qmtpd
bin/qmail-qread
bin/qmail-qstat
bin/qmail-queue
%%LDAP%%bin/qmail-quotawarn
%%SPAMCONTROL%%bin/qmail-recipients
bin/qmail-remote
%%LDAP%%bin/qmail-reply
bin/qmail-rspawn
%%LDAP%%bin/qmail-secretary
bin/qmail-send
bin/qmail-showctl
bin/qmail-smtpd
bin/qmail-start
bin/qmail-tcpok
bin/qmail-tcpto
%%LDAP%%%%EXTERNAL_TODO%%bin/qmail-todo
%%LDAP%%bin/qmail-verify
bin/qreceipt
bin/qsmhook
bin/sendmail
bin/splogger
bin/tcp-env
%%NOT_LDAP%%boot/binm1
%%NOT_LDAP%%boot/binm1+df
%%NOT_LDAP%%boot/binm2
%%NOT_LDAP%%boot/binm2+df
%%NOT_LDAP%%boot/binm3
%%NOT_LDAP%%boot/binm3+df
%%NOT_LDAP%%boot/home
%%NOT_LDAP%%boot/home+df
%%NOT_LDAP%%boot/maildir
%%NOT_LDAP%%boot/proc
%%NOT_LDAP%%boot/proc+df
%%LDAP%%boot/qmail/log/run
%%LDAP%%boot/qmail/run
%%LDAP%%boot/qmail-smtpd/log/run
%%LDAP%%boot/qmail-smtpd/run
%%LDAP%%boot/qmail-qmqpd/log/run
%%LDAP%%boot/qmail-qmqpd/run
%%LDAP%%boot/qmail-pop3d/log/run
%%LDAP%%boot/qmail-pop3d/run
%%LDAP%%boot/qmail-imapd/log/run
%%LDAP%%boot/qmail-imapd/run
%%LDAP%%boot/qmail-pop3d-ssl/log/run
%%LDAP%%boot/qmail-pop3d-ssl/run
%%LDAP%%boot/qmail-imapd-ssl/log/run
%%LDAP%%boot/qmail-imapd-ssl/run
%%LDAP%%boot/qmail-pbsdbd/log/run
%%LDAP%%boot/qmail-pbsdbd/run
@comment This port installs the following manually (ie. not through qmail-hier)
configure/config
configure/config-fast
configure/dnsfq
configure/dnsip
configure/dnsptr
configure/hostname
configure/install
configure/ipmeprint
scripts/enable-qmail
scripts/mkaliasdir
%%SPAMCONTROL%%scripts/qmail-alias2recipients
%%SPAMCONTROL%%scripts/qmail-pwd2recipients
%%SPAMCONTROL%%scripts/qmail-users2recipients
%%SPAMCONTROL%%scripts/qmail-vpopmail2recipients
%%PORTDOCS%%%%DOCSDIR%%/BLURB
%%PORTDOCS%%%%DOCSDIR%%/BLURB2
%%PORTDOCS%%%%DOCSDIR%%/BLURB3
%%PORTDOCS%%%%DOCSDIR%%/BLURB4
%%LDAP%%%%PORTDOCS%%%%DOCSDIR%%/EXTTODO
%%PORTDOCS%%%%DOCSDIR%%/FAQ
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/FILES.spamcontrol
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/HISTORY.spamcontrol
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/INSTALL.alias
%%PORTDOCS%%%%DOCSDIR%%/INSTALL.ctl
%%PORTDOCS%%%%DOCSDIR%%/INSTALL.ids
%%PORTDOCS%%%%DOCSDIR%%/INSTALL.maildir
%%PORTDOCS%%%%DOCSDIR%%/INSTALL.mbox
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/INSTALL.spamcontrol
%%PORTDOCS%%%%DOCSDIR%%/INSTALL.vsm
%%PORTDOCS%%%%DOCSDIR%%/INTERNALS
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/LICENSE.spamcontrol
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/LOGGING.spamcontrol
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/Makefile.djbdns
%%PORTDOCS%%%%DOCSDIR%%/PIC.local2alias
%%PORTDOCS%%%%DOCSDIR%%/PIC.local2ext
%%PORTDOCS%%%%DOCSDIR%%/PIC.local2local
%%PORTDOCS%%%%DOCSDIR%%/PIC.local2rem
%%PORTDOCS%%%%DOCSDIR%%/PIC.local2virt
%%PORTDOCS%%%%DOCSDIR%%/PIC.nullclient
%%PORTDOCS%%%%DOCSDIR%%/PIC.relaybad
%%PORTDOCS%%%%DOCSDIR%%/PIC.relaygood
%%PORTDOCS%%%%DOCSDIR%%/PIC.rem2local
%%LDAP%%%%PORTDOCS%%%%DOCSDIR%%/POPBEFORESMTP
%%PORTDOCS%%%%DOCSDIR%%/PORT_NOTES
%%PORTDOCS%%%%DOCSDIR%%/PORT_NOTES_FreeBSD_40-RELEASE
%%LDAP%%%%PORTDOCS%%%%DOCSDIR%%/QLDAPGROUP
%%LDAP%%%%PORTDOCS%%%%DOCSDIR%%/QLDAPINSTALL
%%LDAP%%%%PORTDOCS%%%%DOCSDIR%%/QLDAPNEWS
%%LDAP%%%%PORTDOCS%%%%DOCSDIR%%/QLDAPPICTURE
%%LDAP%%%%PORTDOCS%%%%DOCSDIR%%/QLDAPTODO
%%README_AUTH%%%%PORTDOCS%%%%DOCSDIR%%/README.auth
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/README.bigtodo
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/README.bouncemaxbytes
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/README.djbdns
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/README.doublebouncetrim
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/README.moreipme
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/README.qmailqueue
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/README.recipients
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/README.spamcontrol
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/README.wildmat
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/README_spamcontrol.html
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/RELEASE_22.spamcontrol
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/RELEASE_23.spamcontrol
%%PORTDOCS%%%%DOCSDIR%%/REMOVE.binmail
%%PORTDOCS%%%%DOCSDIR%%/REMOVE.sendmail
%%PORTDOCS%%%%DOCSDIR%%/SECURITY
%%PORTDOCS%%%%DOCSDIR%%/SENDMAIL
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/SMTPREPLY.spamcontrol
%%SMTP_AUTH_TLS%%%%PORTDOCS%%%%DOCSDIR%%/SMTP_AUTH+TLS.readme
%%PORTDOCS%%%%DOCSDIR%%/SYSDEPS
%%PORTDOCS%%%%DOCSDIR%%/TEST.deliver
%%PORTDOCS%%%%DOCSDIR%%/TEST.receive
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/TESTING.spamcontrol
%%PORTDOCS%%%%DOCSDIR%%/THOUGHTS
%%TLS%%%%PORTDOCS%%%%DOCSDIR%%/TLS.readme
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/TODO.spamcontrol
%%PORTDOCS%%%%DOCSDIR%%/UPGRADE
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/badloadertypes
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/badmailfrom
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/badmimetypes
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/badrcptto
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/conf-spamcontrol
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/install_spamcontrol.sh
%%PORTDOCS%%%%DOCSDIR%%/mailer.conf.sample
%%LDAP%%%%PORTDOCS%%%%DOCSDIR%%/qmail.schema
%%SPAMCONTROL%%%%PORTDOCS%%%%DOCSDIR%%/tarpitcount
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm bin
%%LDAP%%@dirrm boot/qmail/env
%%LDAP%%@dirrm boot/qmail/log
%%LDAP%%@dirrm boot/qmail-smtpd/env
%%LDAP%%@dirrm boot/qmail-smtpd/log
%%LDAP%%@dirrm boot/qmail-smtpd
%%LDAP%%@dirrm boot/qmail-qmqpd/env
%%LDAP%%@dirrm boot/qmail-qmqpd/log
%%LDAP%%@dirrm boot/qmail-qmqpd
%%LDAP%%@dirrm boot/qmail-pop3d/env
%%LDAP%%@dirrm boot/qmail-pop3d/log
%%LDAP%%@dirrm boot/qmail-pop3d
%%LDAP%%@dirrm boot/qmail-imapd/env
%%LDAP%%@dirrm boot/qmail-imapd/log
%%LDAP%%@dirrm boot/qmail-imapd
%%LDAP%%@dirrm boot/qmail-pop3d-ssl/env
%%LDAP%%@dirrm boot/qmail-pop3d-ssl/log
%%LDAP%%@dirrm boot/qmail-pop3d-ssl
%%LDAP%%@dirrm boot/qmail-imapd-ssl/env
%%LDAP%%@dirrm boot/qmail-imapd-ssl/log
%%LDAP%%@dirrm boot/qmail-imapd-ssl
%%LDAP%%@dirrm boot/qmail-pbsdbd/env
%%LDAP%%@dirrm boot/qmail-pbsdbd/log
%%LDAP%%@dirrm boot/qmail-pbsdbd
%%LDAP%%@dirrm boot/qmail
%%LDAP%%@dirrm log/qmail
%%LDAP%%@dirrm log/qmail-smtpd
%%LDAP%%@dirrm log/qmail-qmqpd
%%LDAP%%@dirrm log/qmail-pop3d
%%LDAP%%@dirrm log/qmail-imapd
%%LDAP%%@dirrm log/qmail-pop3d-ssl
%%LDAP%%@dirrm log/qmail-imapd-ssl
%%LDAP%%@dirrm log/qmail-pbsdbd
%%LDAP%%@dirrm log
@dirrm boot
@dirrm configure
@dirrm scripts
@unexec echo "1) If you are permanently removing qmail, you should also ``rm -Rf ${PKG_PREFIX}/alias ${PKG_PREFIX}/control ${PKG_PREFIX}/queue ${PKG_PREFIX}/users'' and remove the qmail users/groups." | /usr/bin/fmt
@unexec echo "2) Do not forget to reconfigure your /etc/mail/mailer.conf to use your older MTA, or reconfigure your old MTA binaries if you decided to replace them." | /usr/bin/fmt
@unexec echo "3) Do not forget to enable sendmail as your default MTA if you have to. You can do it just executing 'echo sendmail_enable=\"YES\" >> /etc/rc.conf" | /usr/bin/fmt
