# comms/atslog/files/atslogd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: atslogd
# REQUIRE: mysql postgresql
# KEYWORD: shutdown 

#
# Add the following line to /etc/rc.conf[.local] to enable atslogd
#
# atslogd_enable (bool):	Set to "NO" by default.
#				Set it to "YES" to enable atslogd.

. /etc/rc.subr

name="atslogd"
rcvar=${name}_enable

load_rc_config $name

: ${atslogd_enable="NO"}

command="%%PREFIX%%/bin/atslogmaster"
required_files="%%PREFIX%%/etc/atslog.conf"
extra_commands="reload writedb rotate alltodb cleardb"
start_cmd="atslogd_command start"
stop_cmd="atslogd_command stop"
reload_cmd="atslogd_command reload"
restart_cmd="atslogd_command restart"
status_cmd="atslogd_command status"
writedb_cmd="atslogd_command writedb"
rotate_cmd="atslogd_command rotate"
alltodb_cmd="atslogd_command alltodb"
cleardb_cmd="atslogd_command cleardb"

atslogd_command()
{
    ${command} ${command_args} ${rc_arg}
}
    
run_rc_command "$1"
