# sysutils/rmonitor/pkg-plist
bin/rmonitor
etc/rc.d/rmonitor.sh
man/man1/rmonitor.1.gz
