# security/stunnel/pkg-plist
etc/rc.d/stunnel.sh
etc/stunnel/stunnel.conf-sample
lib/libstunnel.so
lib/libstunnel.a
sbin/stunnel
sbin/stunnel3
%%PORTDOCS%%%%DOCSDIR%%/stunnel.html
%%PORTDOCS%%%%DOCSDIR%%/stunnel.fr.html
%%PORTDOCS%%%%DOCSDIR%%/stunnel.pl.html
%%PORTDOCS%%%%DOCSDIR%%/examples/ca.html
%%PORTDOCS%%%%DOCSDIR%%/examples/ca.pl
%%PORTDOCS%%%%DOCSDIR%%/examples/importCA.html
%%PORTDOCS%%%%DOCSDIR%%/examples/importCA.sh
%%PORTDOCS%%%%DOCSDIR%%/examples/script.sh
%%PORTDOCS%%%%DOCSDIR%%/examples/stunnel.spec
%%PORTDOCS%%%%DOCSDIR%%/examples/stunnel.init
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%%%DOCSDIR%%/BUGS
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT.GPL
%%PORTDOCS%%%%DOCSDIR%%/CREDITS
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/INSTALL.W32
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/PORTS
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%@dirrm %%DOCSDIR%%/examples
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@unexec rmdir %D/etc/stunnel 2>/dev/null || echo "If you are permanently removing this port, you should do a ``rm -rf ${PKG_PREFIX}/etc/stunnel`` to remove any configuration files left." | fmt
