# multimedia/motion/files/openrc-motion.in
#!/sbin/openrc-run

name=motion

command=%%PREFIX%%/bin/${name}
pidfile="/var/run/${name}.pid"
required_files=%%ETCDIR%%/${name}.conf

depend()
{
	keyword -shutdown
}
