# net-p2p/btsync/pkg-plist
@unexec %D/etc/rc.d/%%PORTNAME%% forcestop 2>/dev/null || true
bin/btsync
@sample etc/%%PORTNAME%%.conf.sample
@unexec rmdir /var/db/%%PORTNAME%% 2>/dev/null || true
@exec install -d -o %%USER%% -g %%GROUP%% /var/db/%%PORTNAME%%
