# audio/jack/files/openrc-jackd.in
#!/sbin/openrc-run

# Note that %%...%% are substituted by the FreeBSD ports
# system
# (see https://www.freebsd.org/doc/en/books/porters-handbook/rc-scripts.html)

# Initialize variables if they are not already set.
: ${jackd_user="root"}
: ${jackd_args="-r -doss -r44100 -p1024 -n3 -w16"}

# These assignments are sufficient for openrc to build a service.
# See: `man openrc-run`

name=jackd
command=%%PREFIX%%/bin/jackd
pidfile=/var/run/jackd.pid
command_args=${jackd_args}
command_usr=${jackd_user}
start_stop_daemon_args="-b -m"

depend()
{
	need localmount
	after network
}
