# mail/milter-callback/files/openrc-milter-callback.in
#!/sbin/openrc-run

name="miltercallback"
command="%%PREFIX%%/libexec/milter-callback"
required_files="%%PREFIX%%/etc/mail/milter-callback.conf"

depend() {
	keyword -shutdown
}
