# comms/conserver/pkg-plist
@comment $FreeBSD$
bin/console
etc/conserver.cf.eg
etc/rc.d/conserver.sh
sbin/autologin
sbin/conserver
