# net-mgmt/arpwatch/files/openrc-arpwatch.in
#!/sbin/openrc-run

arpwatch_dir=${arpwatch_dir:-"%%PREFIX%%/arpwatch/"}

name=arpwatch
required_dirs=${arpwatch_dir}

command="%%PREFIX%%/sbin/${name}"

start_pre() {
    case ${arpwatch_interfaces} in
        '')
	    if [ ! -e "${arpwatch_dir}/arp.dat" ]; then
		if [ -e "${arpwatch_dir}/arp.dat-" ]; then
		    cp "${arpwatch_dir}/arp.dat-" "${arpwatch_dir}/arp.dat"
		else
		    touch "${arpwatch_dir}/arp.dat"
		fi
	    fi
	;;

	*)
	    for interface in ${arpwatch_interfaces}; do
	        if [ ! -e "${arpwatch_dir}/arp.${interface}.dat" ]; then
		    if [ -e "${arpwatch_dir}/arp.${interface}.dat-" ]; then
		        cp "${arpwatch_dir}/arp.${interface}.dat-" "${arpwatch_dir}/arp.${interface}.dat"
		    else
		        touch "${arpwatch_dir}/arp.${interface}.dat"
		    fi
		fi
	    done
	;;
    esac

	if [ -z ${arpwatch_interfaces} ]; then
		 echo 'meuh'
		 #Run as normal
	fi
}

stop() {
    killall arpwatch
}

start() {
	if [ -n ${arpwatch_interfaces} ]; then
		 for interface in ${arpwatch_interfaces}; do
			 rc_interface=`echo ${interface} | sed 's/\./_/g'`
			 eval options=\$arpwatch_${rc_interface}_options
			 command_args="-i ${interface} ${options} -f ${arpwatch_dir}/arp.${interface}.dat"
			 pidfile="/var/run/arpwatch-${interface}.pid"
			 start-stop-daemon --name $rc_interface --pidfile $pidfile --exec $command $command_args
		 done
	fi
}

depend() {
	use net
}
