# security/antivir-milter/pkg-plist
@unexec if cmp -s %D/AntiVir/antivir-dist_avmilter %D/AntiVir/antivir; then rm -f %D/AntiVir/antivir; fi
AntiVir/antivir-dist_avmilter
@exec [ -f %B/antivir ] || cp -p %B/%f %B/antivir
@unexec if cmp -s %D/AntiVir/antivir.vdf-dist_avmilter %D/AntiVir/antivir.vdf; then rm -f %D/AntiVir/antivir.vdf; fi
AntiVir/antivir.vdf-dist_avmilter
@exec [ -f %B/antivir.vdf ] || cp -p %B/%f %B/antivir.vdf
@unexec rmdir %D/AntiVir 2> /dev/null || echo "If permanently deleting this package, %D/AntiVir and its contents must be removed manually."
bin/avq
bin/avqrm
@unexec if cmp -s %D/etc/avmilter/avmilter.conf.sample %D/etc/avmilter/avmilter.conf; then rm -f %D/etc/avmilter/avmilter.conf; else echo "If permanently deleting this package, %D/etc/avmilter/avmilter.conf must be removed manually."; fi
etc/avmilter/avmilter.conf.sample
@exec [ -f %B/avmilter.conf ] || cp %B/%f %B/avmilter.conf
@unexec if cmp -s %D/etc/avmilter/avmilter.ignore.sample %D/etc/avmilter/avmilter.ignore; then rm -f %D/etc/avmilter/avmilter.ignore; else echo "If permanently deleting this package, %D/etc/avmilter/avmilter.ignore must be removed manually."; fi
etc/avmilter/avmilter.ignore.sample
@exec [ -f %B/avmilter.ignore ] || cp %B/%f %B/avmilter.ignore
@unexec if cmp -s %D/etc/avmilter/avmilter.scan.sample %D/etc/avmilter/avmilter.scan; then rm -f %D/etc/avmilter/avmilter.scan; else echo "If permanently deleting this package, %D/etc/avmilter/avmilter.scan must be removed manually."; fi
etc/avmilter/avmilter.scan.sample
@exec [ -f %B/avmilter.scan ] || cp %B/%f %B/avmilter.scan
@unexec if cmp -s %D/etc/avmilter/avmilter.warn.sample %D/etc/avmilter/avmilter.warn; then rm -f %D/etc/avmilter/avmilter.warn; else echo "If permanently deleting this package, %D/etc/avmilter/avmilter.warn must be removed manually."; fi
etc/avmilter/avmilter.warn.sample
@exec [ -f %B/avmilter.warn ] || cp %B/%f %B/avmilter.warn
@unexec rmdir %D/etc/avmilter 2> /dev/null || echo "If permanently deleting this package, %D/etc/avmilter and its contents must be removed manually."
@unexec if cmp -s %D/etc/antivir.conf.sample %D/etc/antivir.conf; then rm -f %D/etc/antivir.conf; else echo "If permanently deleting this package, %D/etc/antivir.conf must be removed manually."; fi
etc/antivir.conf.sample
@exec [ -f %B/antivir.conf ] || cp %B/%f %B/antivir.conf
@cwd /
etc/antivir.conf
@exec ln -sf %%LOCALBASE%%/%F %D/%F
@cwd %%LOCALBASE%%
etc/rc.d/avmilter.sh
etc/rc.d/avupdater.sh-dist_avmilter
@unexec if [ -f %B/avupdater.sh ]; then echo "If permanently deleting this package, %B/avupdater.sh must be removed manually."; fi
sbin/antivirupdater
sbin/avmilter
@exec install -d -o smmsp -g smmsp -m 700 /var/spool/avmilter
@exec install -d -o smmsp -g smmsp -m 700 /var/spool/avmilter/incoming
@exec install -d -o smmsp -g smmsp -m 700 /var/spool/avmilter/outgoing
@exec install -d -o smmsp -g smmsp -m 700 /var/spool/avmilter/rejected
@unexec rmdir /var/spool/avmilter/incoming 2> /dev/null || true
@unexec rmdir /var/spool/avmilter/outgoing 2> /dev/null || true
@unexec rmdir /var/spool/avmilter/rejected 2> /dev/null || true
@unexec rmdir /var/spool/avmilter 2> /dev/null || echo "If permanently deleting this package, /var/spool/avmilter and its contents must be removed manually."
%%EXAMPLESDIR%%/templates/de/patho-administrator
%%EXAMPLESDIR%%/templates/de/patho-recipient
%%EXAMPLESDIR%%/templates/de/patho-sender
%%EXAMPLESDIR%%/templates/de/virus-administrator
%%EXAMPLESDIR%%/templates/de/virus-recipient
%%EXAMPLESDIR%%/templates/de/virus-sender
@dirrm %%EXAMPLESDIR%%/templates/de
%%EXAMPLESDIR%%/templates/en/patho-administrator
%%EXAMPLESDIR%%/templates/en/patho-recipient
%%EXAMPLESDIR%%/templates/en/patho-sender
%%EXAMPLESDIR%%/templates/en/virus-administrator
%%EXAMPLESDIR%%/templates/en/virus-recipient
%%EXAMPLESDIR%%/templates/en/virus-sender
@dirrm %%EXAMPLESDIR%%/templates/en
%%EXAMPLESDIR%%/templates/es/patho-administrator
%%EXAMPLESDIR%%/templates/es/patho-recipient
%%EXAMPLESDIR%%/templates/es/patho-sender
%%EXAMPLESDIR%%/templates/es/virus-administrator
%%EXAMPLESDIR%%/templates/es/virus-recipient
%%EXAMPLESDIR%%/templates/es/virus-sender
@dirrm %%EXAMPLESDIR%%/templates/es
%%EXAMPLESDIR%%/templates/hu/patho-administrator
%%EXAMPLESDIR%%/templates/hu/patho-recipient
%%EXAMPLESDIR%%/templates/hu/patho-sender
%%EXAMPLESDIR%%/templates/hu/virus-administrator
%%EXAMPLESDIR%%/templates/hu/virus-recipient
%%EXAMPLESDIR%%/templates/hu/virus-sender
@dirrm %%EXAMPLESDIR%%/templates/hu
%%EXAMPLESDIR%%/templates/it/patho-administrator
%%EXAMPLESDIR%%/templates/it/patho-recipient
%%EXAMPLESDIR%%/templates/it/patho-sender
%%EXAMPLESDIR%%/templates/it/virus-administrator
%%EXAMPLESDIR%%/templates/it/virus-recipient
%%EXAMPLESDIR%%/templates/it/virus-sender
@dirrm %%EXAMPLESDIR%%/templates/it
%%EXAMPLESDIR%%/templates/nl/patho-administrator
%%EXAMPLESDIR%%/templates/nl/patho-recipient
%%EXAMPLESDIR%%/templates/nl/patho-sender
%%EXAMPLESDIR%%/templates/nl/virus-administrator
%%EXAMPLESDIR%%/templates/nl/virus-recipient
%%EXAMPLESDIR%%/templates/nl/virus-sender
@dirrm %%EXAMPLESDIR%%/templates/nl
@dirrm %%EXAMPLESDIR%%/templates
@dirrm %%EXAMPLESDIR%%
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/LICENSE
%%PORTDOCS%%%%DOCSDIR%%/LICENSE.DE
%%PORTDOCS%%%%DOCSDIR%%/MANUAL
%%PORTDOCS%%%%DOCSDIR%%/README.gpg
%%PORTDOCS%%%%DOCSDIR%%/avmilter_de.pdf
%%PORTDOCS%%%%DOCSDIR%%/avmilter_en.pdf
%%PORTDOCS%%%%DOCSDIR%%/antivir.gpg
%%PORTDOCS%%@dirrm %%DOCSDIR%%
