# japanese/Wnn6/pkg-plist
OMRONWnn6/doc/copyright.html
OMRONWnn6/doc/man07.0006.i.jpg
OMRONWnn6/doc/man07.0004.i.jpg
OMRONWnn6/doc/man07.0003.i.jpg
OMRONWnn6/doc/man06.html
OMRONWnn6/doc/man05.html
OMRONWnn6/doc/man05.0002.i.jpg
OMRONWnn6/doc/man04.html
OMRONWnn6/doc/man04.0005.i.jpg
OMRONWnn6/doc/man04.0003.i.jpg
OMRONWnn6/doc/man03.html
OMRONWnn6/doc/man03.0018.i.jpg
OMRONWnn6/doc/man03.0016.i.jpg
OMRONWnn6/doc/man03.0006.i.jpg
OMRONWnn6/doc/man02.html
OMRONWnn6/doc/man02.0033.i.jpg
OMRONWnn6/doc/man02.0006.i.jpg
OMRONWnn6/doc/man02.0002.i.jpg
OMRONWnn6/doc/man01.html
OMRONWnn6/doc/intro-02.html
OMRONWnn6/doc/intro-02.0001.i.jpg
OMRONWnn6/doc/intro-01.html
OMRONWnn6/doc/index.html
OMRONWnn6/doc/hajimeni.html
OMRONWnn6/doc/furoku-symbol.html
OMRONWnn6/doc/furoku-error.html
OMRONWnn6/doc/man07.html
OMRONWnn6/wnn6dic/bushu/bushu.dic
OMRONWnn6/wnn6dic/bushu/bushu.fzk
OMRONWnn6/wnn6dic/iwanami/fisd
OMRONWnn6/wnn6dic/iwanami/ikeiji.dic
OMRONWnn6/wnn6dic/iwanami/kihon.dic
OMRONWnn6/wnn6dic/iwanami/kougo.fzk
OMRONWnn6/wnn6dic/iwanami/std.fzk
OMRONWnn6/wnn6dic/iwanami/symbol.dic
OMRONWnn6/wnn6dic/iwanami/tankan.dic
OMRONWnn6/wnn6dic/iwanami/tankan2.dic
OMRONWnn6/wnn6dic/iwanami/tankan3.dic
OMRONWnn6/wnn6dic/iwanami/tel.dic
OMRONWnn6/wnn6dic/iwanami/zip.dic
OMRONWnn6/wnn6fbsd/ja_JP/jserverrc
OMRONWnn6/wnn6fbsd/ja_JP/hinsi.data
OMRONWnn6/wnn6fbsd/ja_JP/wnnstat.msg
OMRONWnn6/wnn6fbsd/ja_JP/libwnn.msg
OMRONWnn6/wnn6fbsd/ja_JP/uumrc
OMRONWnn6/wnn6fbsd/ja_JP/wnnenvrc
OMRONWnn6/wnn6fbsd/ja_JP/otow.format/atok7-wnn6.fmt
OMRONWnn6/wnn6fbsd/ja_JP/otow.format/atok8-wnn6.fmt
OMRONWnn6/wnn6fbsd/ja_JP/otow.format/egbridge-wnn6.fmt
OMRONWnn6/wnn6fbsd/ja_JP/otow.format/vje-wnn6.fmt
OMRONWnn6/wnn6fbsd/offlinerc
OMRONWnn6/wnn6fbsd/serverdefs
OMRONWnn6/wnn6fbsd/udmergerc
OMRONWnn6/wnnsdk/mule-2.2.2/README.mule
OMRONWnn6/wnnsdk/mule-2.2.2/mule-2.2.2-wnn6.diff
OMRONWnn6/wnnsdk/sdk/SunOS-4.1.tar.gz
OMRONWnn6/wnnsdk/sdk/Solaris-2.4.tar.gz
OMRONWnn6/wnnsdk/sdk/IRIX-5.3.tar.gz
OMRONWnn6/wnnsdk/sdk/DGUX-5.4.tar.gz
OMRONWnn6/wnnsdk/sdk/Solaris-2.4-x86.tar.gz
OMRONWnn6/wnnsdk/sdk/BSD386-1.1.tar.gz
OMRONWnn6/wnnsdk/sdk/AIX-3.2.tar.gz
OMRONWnn6/wnnsdk/sdk/AlphaOSF1-1.2.tar.gz
OMRONWnn6/wnnsdk/sdk/EWS-UX-4.0.tar.gz
OMRONWnn6/wnnsdk/sdk/FreeBSD-2.0.5R.tar.gz
OMRONWnn6/wnnsdk/sdk/HI-UX-3050RX.tar.gz
OMRONWnn6/wnnsdk/sdk/HPUX-9.0.tar.gz
OMRONWnn6/wnnsdk/sdk/LUNA88K-2.0.tar.gz
OMRONWnn6/wnnsdk/sdk/Linux-1.2.8.tar.gz
OMRONWnn6/wnnsdk/sdk/NEWS-OS-4.2R.tar.gz
OMRONWnn6/wnnsdk/sdk/NEWS-OS-6.0.tar.gz
OMRONWnn6/wnnsdk/sdk/UniOS-U-3.0.tar.gz
OMRONWnn6/wnnsdk/sdk/AlphaOSF1-3.2A.tar.gz
OMRONWnn6/wnnsdk/sdk/FreeBSD-2.2R.tar.gz
OMRONWnn6/wnnsdk/sdk/Linux-2.0.0.tar.gz
OMRONWnn6/Wnn6fbsd/atod
OMRONWnn6/Wnn6fbsd/dtoa
OMRONWnn6/Wnn6fbsd/jserver
OMRONWnn6/Wnn6fbsd/otow
OMRONWnn6/Wnn6fbsd/udmerge
OMRONWnn6/Wnn6fbsd/wnnaccess
OMRONWnn6/Wnn6fbsd/wnnds
OMRONWnn6/Wnn6fbsd/wnnkill
OMRONWnn6/Wnn6fbsd/wnnoffline
OMRONWnn6/Wnn6fbsd/wnnstat
OMRONWnn6/Wnn6fbsd/wnntouch
OMRONWnn6/README.sjis
OMRONWnn6/README.jis
OMRONWnn6/README.euc
@dirrm OMRONWnn6/wnnsdk/sdk
@dirrm OMRONWnn6/wnnsdk/mule-2.2.2
@dirrm OMRONWnn6/wnnsdk
@dirrm OMRONWnn6/wnn6fbsd/ja_JP/otow.format
@unexec rm %D/OMRONWnn6/wnn6fbsd/ja_JP/dic
@dirrm OMRONWnn6/wnn6fbsd/ja_JP
@dirrm OMRONWnn6/wnn6fbsd
@dirrm OMRONWnn6/wnn6dic/iwanami
@dirrm OMRONWnn6/wnn6dic/bushu
@dirrm OMRONWnn6/wnn6dic
@dirrm OMRONWnn6/doc
@dirrm OMRONWnn6/Wnn6fbsd
@dirrm OMRONWnn6
etc/rc.d/Wnn6.sh
@unexec rm /usr/local/bin/Wnn6
@unexec rm /usr/local/lib/wnn6
