# emulators/virtualbox-ose-additions/files/openrc-vboxservice.in
#!/sbin/openrc-run

name="vboxservice"
command="%%PREFIX%%/sbin/VBoxService"

depend() {
	need vboxguest
	keyword -shutdown
}
