# emulators/virtualbox-ose-additions/files/openrc-vboxguest.in
#!/sbin/openrc-run

name="vboxguest"
vboxguest_modules="vboxguest:pci/vboxguest"

start()
{
	local _k _m

	for _m in ${vboxguest_modules}; do
		_k=${_m%:*}
		_m=${_m#*:}
		if ! kldstat -q -m ${_m}; then
			if ! kldload ${_k} > /dev/null 2>&1; then
				ewarn "Can't load ${_k} module."
				return 1
			fi
		fi
	done
}

stop()
{
	local _k _m _r

	_r=
	for _m in ${vboxguest_modules}; do
		_r="${_m} ${_r}"
	done
	for _m in ${_r}; do
		_k=${_m%:*}
		_m=${_m#*:}
		if kldstat -q -m ${_m}; then
			if ! kldunload ${_k} > /dev/null 2>&1; then
				ewarn "Can't unload ${_k} module."
				return 1
			fi
		fi
	done
}

depend() {
	provide vboxguest
	keyword -jail
}
