# sysutils/yum/pkg-plist
bin/yum
man/man5/yum-updatesd.conf.5.gz
man/man5/yum.conf.5.gz
man/man8/yum-shell.8.gz
man/man8/yum-updatesd.8.gz
man/man8/yum.8.gz
sbin/yum-updatesd
share/yum-cli/callback.py
share/yum-cli/callback.pyc
share/yum-cli/cli.py
share/yum-cli/cli.pyc
share/yum-cli/output.py
share/yum-cli/output.pyc
share/yum-cli/shell.py
share/yum-cli/shell.pyc
share/yum-cli/utils.py
share/yum-cli/utils.pyc
share/yum-cli/yumcommands.py
share/yum-cli/yumcommands.pyc
share/yum-cli/yummain.py
share/yum-cli/yummain.pyc
share/yum-cli/yumupd.py
share/yum-cli/yumupd.pyc
%%PYTHON_SITELIBDIR%%/rpmUtils/__init__.py
%%PYTHON_SITELIBDIR%%/rpmUtils/__init__.pyc
%%PYTHON_SITELIBDIR%%/rpmUtils/arch.py
%%PYTHON_SITELIBDIR%%/rpmUtils/arch.pyc
%%PYTHON_SITELIBDIR%%/rpmUtils/miscutils.py
%%PYTHON_SITELIBDIR%%/rpmUtils/miscutils.pyc
%%PYTHON_SITELIBDIR%%/rpmUtils/oldUtils.py
%%PYTHON_SITELIBDIR%%/rpmUtils/oldUtils.pyc
%%PYTHON_SITELIBDIR%%/rpmUtils/transaction.py
%%PYTHON_SITELIBDIR%%/rpmUtils/transaction.pyc
%%PYTHON_SITELIBDIR%%/rpmUtils/updates.py
%%PYTHON_SITELIBDIR%%/rpmUtils/updates.pyc
%%PYTHON_SITELIBDIR%%/yum/Errors.py
%%PYTHON_SITELIBDIR%%/yum/Errors.pyc
%%PYTHON_SITELIBDIR%%/yum/__init__.py
%%PYTHON_SITELIBDIR%%/yum/__init__.pyc
%%PYTHON_SITELIBDIR%%/yum/callbacks.py
%%PYTHON_SITELIBDIR%%/yum/callbacks.pyc
%%PYTHON_SITELIBDIR%%/yum/comps.py
%%PYTHON_SITELIBDIR%%/yum/comps.pyc
%%PYTHON_SITELIBDIR%%/yum/config.py
%%PYTHON_SITELIBDIR%%/yum/config.pyc
%%PYTHON_SITELIBDIR%%/yum/constants.py
%%PYTHON_SITELIBDIR%%/yum/constants.pyc
%%PYTHON_SITELIBDIR%%/yum/depsolve.py
%%PYTHON_SITELIBDIR%%/yum/depsolve.pyc
%%PYTHON_SITELIBDIR%%/yum/failover.py
%%PYTHON_SITELIBDIR%%/yum/failover.pyc
%%PYTHON_SITELIBDIR%%/yum/history.py
%%PYTHON_SITELIBDIR%%/yum/history.pyc
%%PYTHON_SITELIBDIR%%/yum/i18n.py
%%PYTHON_SITELIBDIR%%/yum/i18n.pyc
%%PYTHON_SITELIBDIR%%/yum/logginglevels.py
%%PYTHON_SITELIBDIR%%/yum/logginglevels.pyc
%%PYTHON_SITELIBDIR%%/yum/mdparser.py
%%PYTHON_SITELIBDIR%%/yum/mdparser.pyc
%%PYTHON_SITELIBDIR%%/yum/metalink.py
%%PYTHON_SITELIBDIR%%/yum/metalink.pyc
%%PYTHON_SITELIBDIR%%/yum/misc.py
%%PYTHON_SITELIBDIR%%/yum/misc.pyc
%%PYTHON_SITELIBDIR%%/yum/packageSack.py
%%PYTHON_SITELIBDIR%%/yum/packageSack.pyc
%%PYTHON_SITELIBDIR%%/yum/packages.py
%%PYTHON_SITELIBDIR%%/yum/packages.pyc
%%PYTHON_SITELIBDIR%%/yum/parser.py
%%PYTHON_SITELIBDIR%%/yum/parser.pyc
%%PYTHON_SITELIBDIR%%/yum/pgpmsg.py
%%PYTHON_SITELIBDIR%%/yum/pgpmsg.pyc
%%PYTHON_SITELIBDIR%%/yum/pkgtag_db.py
%%PYTHON_SITELIBDIR%%/yum/pkgtag_db.pyc
%%PYTHON_SITELIBDIR%%/yum/plugins.py
%%PYTHON_SITELIBDIR%%/yum/plugins.pyc
%%PYTHON_SITELIBDIR%%/yum/repoMDObject.py
%%PYTHON_SITELIBDIR%%/yum/repoMDObject.pyc
%%PYTHON_SITELIBDIR%%/yum/repos.py
%%PYTHON_SITELIBDIR%%/yum/repos.pyc
%%PYTHON_SITELIBDIR%%/yum/rpmsack.py
%%PYTHON_SITELIBDIR%%/yum/rpmsack.pyc
%%PYTHON_SITELIBDIR%%/yum/rpmtrans.py
%%PYTHON_SITELIBDIR%%/yum/rpmtrans.pyc
%%PYTHON_SITELIBDIR%%/yum/sqlitesack.py
%%PYTHON_SITELIBDIR%%/yum/sqlitesack.pyc
%%PYTHON_SITELIBDIR%%/yum/sqlutils.py
%%PYTHON_SITELIBDIR%%/yum/sqlutils.pyc
%%PYTHON_SITELIBDIR%%/yum/transactioninfo.py
%%PYTHON_SITELIBDIR%%/yum/transactioninfo.pyc
%%PYTHON_SITELIBDIR%%/yum/update_md.py
%%PYTHON_SITELIBDIR%%/yum/update_md.pyc
%%PYTHON_SITELIBDIR%%/yum/yumRepo.py
%%PYTHON_SITELIBDIR%%/yum/yumRepo.pyc
share/locale/bn_IN/LC_MESSAGES/yum.mo
share/locale/ca/LC_MESSAGES/yum.mo
share/locale/cs/LC_MESSAGES/yum.mo
share/locale/da/LC_MESSAGES/yum.mo
share/locale/de/LC_MESSAGES/yum.mo
share/locale/el/LC_MESSAGES/yum.mo
share/locale/en_GB/LC_MESSAGES/yum.mo
share/locale/es/LC_MESSAGES/yum.mo
share/locale/eu/LC_MESSAGES/yum.mo
share/locale/fi/LC_MESSAGES/yum.mo
share/locale/fr/LC_MESSAGES/yum.mo
share/locale/gu/LC_MESSAGES/yum.mo
share/locale/hi/LC_MESSAGES/yum.mo
share/locale/hu/LC_MESSAGES/yum.mo
share/locale/id/LC_MESSAGES/yum.mo
share/locale/id_ID/LC_MESSAGES/yum.mo
share/locale/it/LC_MESSAGES/yum.mo
share/locale/ja/LC_MESSAGES/yum.mo
share/locale/lt_LT/LC_MESSAGES/yum.mo
share/locale/mr/LC_MESSAGES/yum.mo
share/locale/ms/LC_MESSAGES/yum.mo
share/locale/nb/LC_MESSAGES/yum.mo
share/locale/pa/LC_MESSAGES/yum.mo
share/locale/pl/LC_MESSAGES/yum.mo
share/locale/pt/LC_MESSAGES/yum.mo
share/locale/pt_BR/LC_MESSAGES/yum.mo
share/locale/ru/LC_MESSAGES/yum.mo
share/locale/sr/LC_MESSAGES/yum.mo
share/locale/sr@latin/LC_MESSAGES/yum.mo
share/locale/sv/LC_MESSAGES/yum.mo
share/locale/uk/LC_MESSAGES/yum.mo
share/locale/ur/LC_MESSAGES/yum.mo
share/locale/zh_CN/LC_MESSAGES/yum.mo
share/locale/zh_TW/LC_MESSAGES/yum.mo
etc/logrotate.d/yum
etc/rc.d/init.d/yum-updatesd
etc/dbus-1/system.d/yum-updatesd.conf
etc/bash_completion.d/yum.bash
etc/yum/version-groups.conf
@unexec if cmp -s %D/etc/yum/yum.conf.sample %D/etc/yum/yum.conf; then rm -f %D/etc/yum/yum.conf; fi
etc/yum/yum.conf.sample
@exec if [ ! -f %D/etc/yum/yum.conf ]; then cp -p %D/%F %B/yum.conf; fi
@unexec if cmp -s %D/etc/yum/yum-updatesd.conf.sample %D/etc/yum/yum-updatesd.conf; then rm -f %D/etc/yum/yum-updatesd.conf; fi
etc/yum/yum-updatesd.conf.sample
@exec if [ ! -f %D/etc/yum/yum-updatesd.conf ]; then cp -p %D/%F %B/yum-updatesd.conf; fi
@dirrm share/yum-cli
@dirrm etc/yum/repos.d
@dirrm etc/yum/vars
@dirrm etc/yum/protected.d
@dirrm etc/yum
@dirrm %%PYTHON_SITELIBDIR%%/rpmUtils
@dirrm %%PYTHON_SITELIBDIR%%/yum
@dirrmtry etc/bash_completion.d
@dirrmtry etc/logrotate.d
@dirrmtry etc/rc.d/init.d
@dirrmtry share/locale/id_ID/LC_MESSAGES
@dirrmtry share/locale/id_ID
@dirrmtry share/locale/lt_LT/LC_MESSAGES
@dirrmtry share/locale/lt_LT
@dirrmtry share/locale/ur/LC_MESSAGES
@dirrmtry share/locale/ur
@dirrmtry etc/dbus-1/system.d
@dirrmtry etc/dbus-1
@unexec rmdir /var/cache/yum 2>&1 >/dev/null || true
@unexec rmdir /var/lib/yum 2>&1 >/dev/null || true
@unexec rmdir /var/lib 2>&1 >/dev/null || true
