# databases/couchdb/pkg-plist
bin/couchdb
bin/couchjs
%%ETCDIR%%/conf/mime.types
%%ETCDIR%%/couch.ini
%%ETCDIR%%/couch_httpd.conf
etc/default/couchdb
etc/rc.d/couchdb
lib/couchdb/bin/_couchjs
lib/couchdb/erlang/bin/couch.boot
lib/couchdb/erlang/lib/couch-0.7.0/ebin/cjson.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch.app
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_btree.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_db.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_db_sup.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_db_update_notifier.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_doc.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_event_sup.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_file.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_ft_query.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_log.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_query_servers.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_rep.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_server.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_server_sup.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_stream.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_util.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/couch_view_group.beam
lib/couchdb/erlang/lib/couch-0.7.0/ebin/mod_couch.beam
lib/couchdb/erlang/lib/couch-0.7.0/include/couch_db.hrl
lib/couchdb/erlang/lib/couch-0.7.0/priv/lib/couch_erl_driver.a
lib/couchdb/erlang/lib/couch-0.7.0/priv/lib/couch_erl_driver.la
lib/couchdb/erlang/lib/couch-0.7.0/priv/lib/couch_erl_driver.so
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/couch_inets.app
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/ftp.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/ftp_progress.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/ftp_response.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/ftp_sup.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/http.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/http_base_64.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/http_chunk.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/http_cookie.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/http_request.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/http_response.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/http_transport.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/http_uri.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/http_util.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpc_handler.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpc_manager.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpc_request.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpc_response.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpc_sup.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_acceptor.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_acceptor_sup.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_cgi.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_conf.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_esi.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_example.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_instance_sup.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_manager.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_misc_sup.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_request.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_request_handler.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_response.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_script_env.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_socket.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_sup.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/httpd_util.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/inets.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/inets_app.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/inets_sup.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_actions.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_alias.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_auth.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_auth_dets.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_auth_mnesia.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_auth_plain.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_auth_server.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_browser.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_cgi.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_dir.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_disk_log.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_esi.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_get.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_head.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_htaccess.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_include.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_log.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_range.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_responsecontrol.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_security.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_security_server.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/mod_trace.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/tftp.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/tftp_binary.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/tftp_engine.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/tftp_file.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/tftp_lib.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin/tftp_sup.beam
lib/couchdb/erlang/lib/couch_inets-4.7.5/include/ftp_internal.hrl
lib/couchdb/erlang/lib/couch_inets-4.7.5/include/http_internal.hrl
lib/couchdb/erlang/lib/couch_inets-4.7.5/include/httpc_internal.hrl
lib/couchdb/erlang/lib/couch_inets-4.7.5/include/httpd.hrl
lib/couchdb/erlang/lib/couch_inets-4.7.5/include/inets_internal.hrl
lib/couchdb/erlang/lib/couch_inets-4.7.5/include/mod_auth.hrl
lib/couchdb/erlang/lib/couch_inets-4.7.5/include/tftp.hrl
%%DATADIR%%/server/main.js
%%DATADIR%%/www/browse/database.html
%%DATADIR%%/www/browse/document.html
%%DATADIR%%/www/browse/index.html
%%DATADIR%%/www/couch_tests.html
%%DATADIR%%/www/favicon.ico
%%DATADIR%%/www/image/add.gif
%%DATADIR%%/www/image/apply.gif
%%DATADIR%%/www/image/bg.png
%%DATADIR%%/www/image/cancel.gif
%%DATADIR%%/www/image/delete-mini.gif
%%DATADIR%%/www/image/delete.gif
%%DATADIR%%/www/image/hgrad.gif
%%DATADIR%%/www/image/load.gif
%%DATADIR%%/www/image/logo.png
%%DATADIR%%/www/image/path.gif
%%DATADIR%%/www/image/run-mini.gif
%%DATADIR%%/www/image/run.gif
%%DATADIR%%/www/image/running.gif
%%DATADIR%%/www/image/save.gif
%%DATADIR%%/www/image/test_failure.gif
%%DATADIR%%/www/image/test_success.gif
%%DATADIR%%/www/image/thead.gif
%%DATADIR%%/www/image/toggle-collapse.gif
%%DATADIR%%/www/image/toggle-expand.gif
%%DATADIR%%/www/index.html
%%DATADIR%%/www/replicator.html
%%DATADIR%%/www/script/browse.js
%%DATADIR%%/www/script/couch.js
%%DATADIR%%/www/script/couch_tests.js
%%DATADIR%%/www/script/jquery.js
%%DATADIR%%/www/script/json.js
%%DATADIR%%/www/script/pprint.js
%%DATADIR%%/www/script/shell.js
%%DATADIR%%/www/shell.html
%%DATADIR%%/www/style/layout.css
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS.gz
%%PORTDOCS%%%%DOCSDIR%%/BUGS.gz
%%PORTDOCS%%%%DOCSDIR%%/NEWS.gz
%%PORTDOCS%%%%DOCSDIR%%/README.gz
%%PORTDOCS%%%%DOCSDIR%%/THANKS.gz
%%PORTDOCS%%%%DOCSDIR%%/changelog.gz
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@dirrm %%DATADIR%%/www/style
@dirrm %%DATADIR%%/www/script
@dirrm %%DATADIR%%/www/image
@dirrm %%DATADIR%%/www/browse
@dirrm %%DATADIR%%/www
@dirrm %%DATADIR%%/server
@dirrm %%DATADIR%%
@dirrm lib/couchdb/erlang/lib/couch_inets-4.7.5/include
@dirrm lib/couchdb/erlang/lib/couch_inets-4.7.5/ebin
@dirrm lib/couchdb/erlang/lib/couch_inets-4.7.5
@dirrm lib/couchdb/erlang/lib/couch-0.7.0/priv/lib
@dirrm lib/couchdb/erlang/lib/couch-0.7.0/priv
@dirrm lib/couchdb/erlang/lib/couch-0.7.0/include
@dirrm lib/couchdb/erlang/lib/couch-0.7.0/ebin
@dirrm lib/couchdb/erlang/lib/couch-0.7.0
@dirrm lib/couchdb/erlang/lib
@dirrm lib/couchdb/erlang/bin
@dirrm lib/couchdb/erlang
@dirrm lib/couchdb/bin
@dirrm lib/couchdb
@dirrm etc/default
@dirrm %%ETCDIR%%/conf
@dirrm %%ETCDIR%%
