# devel/oozie/files/oozied.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: oozied
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# oozied_enable (bool):          Set to NO by default.
#                                Set it to YES to enable oozied.

. /etc/rc.subr

export PATH=${PATH}:%%LOCALBASE%%/bin
name=oozied
rcvar=oozied_enable
pidfile=%%OOZIE_RUNDIR%%/${name}.pid

load_rc_config "${name}"

: ${oozied_enable:=NO}
oozied_user=%%OOZIE_USER%%

command="%%PREFIX%%/oozie/bin/oozied.sh"
command_args="start"

oozied_stop () {
  ${command} stop
}

run_rc_command "$1"
