# www/dhttpd/files/openrc-dhttpd.in
#!/sbin/openrc-run

name="dhttpd"
description=""

command=%%PREFIX%%/sbin/dhttpd

depend(){
	keyword -shutdown
}
