# net/zaptel12/pkg-plist
bin/fxo_test
bin/ztcfg
bin/ztmonitor
bin/zttest
bin/zttool
etc/rc.d/zaptel.sh.sample
@unexec if cmp -s %D/etc/zaptel.conf %D/etc/zaptel.conf.sample; then rm -f %D/etc/zaptel.conf; fi
etc/zaptel.conf.sample
@exec [ -f %B/zaptel.conf ] || cp %B/%f %B/zaptel.conf
include/tonezone.h
include/zaptel.h
lib/libtonezone.a
lib/zaptel/linker.hints
lib/zaptel/wcfxo.ko
lib/zaptel/wcfxs.ko
lib/zaptel/zaptel.ko
lib/zaptel/ztdummy.ko
@dirrm lib/zaptel
