# mail/perdition/files/openrc-perdition.in
#!/sbin/openrc-run

# perdition_flags="<set as needed>"
#
# See perdition(8) for flags
#

name="perdition"

command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/perdition/${name}.pid"

depend()
{
	keyword -shutdown
}
