# www/junkbuster/files/openrc-junkbuster.in
#!/sbin/openrc-run

name="junkbuster"
description=""
command=%%PREFIX%%/sbin/junkbuster
required_files=%%PREFIX%%/etc/junkbuster/configfile
