# mail/policyd2/files/policyd2.in
#!/bin/sh

# Start or stop policyd2
# $FreeBSD$

# PROVIDE: policyd2
# KEYWORD: shutdown
#
# Define these policyd2_* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#

. /etc/rc.subr

name="policyd2"
rcvar=policyd2_enable

command="%%LOCALBASE%%/bin/cbpolicyd"
command_interpreter="/usr/bin/perl"

load_rc_config $name

: ${policyd2_enable="NO"}

run_rc_command "$1"
