# net/pdnsd/pkg-plist
etc/pdnsd.conf.sample
etc/rc.d/pdnsd.sh
sbin/pdnsd
sbin/pdnsd-ctl
@unexec rm -rf /var/db/pdnsd
