# net-im/py-xmpppy-irc/files/xmpppyirc.in
#!/bin/sh

# PROVIDE: xmpppyirc
# REQUIRE: DAEMON

. /etc/rc.subr

name="xmpppyirc"
pidfile="/var/run/$name.pid"

command="%%FEXEC%%"
command_args="&"
command_interpreter="%%PYTHON_CMD%%"
required_files="%%PREFIX%%/etc/$name.conf"
start_precmd="touch $pidfile && export PID=$pidfile"

load_rc_config $name
: ${xmpppyirc_enable="NO"}
run_rc_command "$1"
