# net-mgmt/xymon-client/files/hobbit-client.sh.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: hobbit_client
# REQUIRE: DAEMON

. /etc/rc.subr

name="hobbit_client"
rcvar=`set_rcvar`

load_rc_config "$name"
: ${hobbit_client_enable:="NO"}
: ${hobbit_client_flags:=""}

command="%%PREFIX%%/www/xymon/client/runclient.sh"
start_cmd="%%PREFIX%%/www/xymon/client/runclient.sh ${hobbit_client_flags} start"
stop_cmd="%%PREFIX%%/www/xymon/client/runclient.sh ${hobbit_client_flags} stop"
status_cmd="%%PREFIX%%/www/xymon/client/runclient.sh ${hobbit_client_flags} status"
reload_cmd="%%PREFIX%%/www/xymon/client/runclient.sh ${hobbit_client_flags} restart"
pidfile="%%PREFIX%%/www/xymon/client/logs/clientlaunch.`hostname`.pid"

run_rc_command "$1"
# eof
