# net-mgmt/xymon-client/files/xymon-client.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: xymon_client
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name=xymon_client
rcvar=xymon_client_enable

load_rc_config "$name"
: ${xymon_client_enable:=NO}
: ${xymon_client_user:=%%XYMONUSER%%}

command=%%WWWDIR%%/client/runclient.sh
command_args="${xymon_client_flags} ${1}"
procname=%%WWWDIR%%/client/bin/xymonlaunch
pidfile="%%WWWDIR%%/client/logs/clientlaunch.`hostname`.pid"
start_precmd="chown -R ${xymon_client_user} %%WWWDIR%%/client/logs"

run_rc_command "$1"
