# korean/FreeWnn-lib/pkg-plist
@cwd /usr/local
etc/rc.d/kWnn.sh
lib/wnn/cvt_key_empty
lib/wnn/cvt_key_tbl
lib/wnn/cvt_key_tbl.ST
lib/wnn/cvt_key_tbl.gm
lib/wnn/cvt_key_tbl.kt
lib/wnn/cvt_key_tbl.mv
lib/wnn/cvt_key_tbl.vt
lib/wnn/ko_KR/dic/sys/full.fzk
lib/wnn/ko_KR/dic/sys/hword.dic
lib/wnn/ko_KR/dic/sys/single.dic
lib/wnn/ko_KR/hinsi.data
lib/wnn/ko_KR/kserverrc
lib/wnn/ko_KR/libwnn.msg
lib/wnn/ko_KR/rk/1B_TOUPPER
lib/wnn/ko_KR/rk/2A_CTRL
lib/wnn/ko_KR/rk/2C_HANGUL.KEY
lib/wnn/ko_KR/rk/2C_HANGUL.ROM
lib/wnn/ko_KR/rk/2C_KSC
lib/wnn/ko_KR/rk/3B_JEONKAK
lib/wnn/ko_KR/rk/mode
lib/wnn/ko_KR/rk/mode.key
lib/wnn/ko_KR/rk/mode.rom
lib/wnn/ko_KR/uum.msg
lib/wnn/ko_KR/uumkey
lib/wnn/ko_KR/uumrc
lib/wnn/ko_KR/wnnenvrc
lib/wnn/ko_KR/wnnenvrc_R
lib/wnn/ko_KR/wnnstat.msg
lib/wnn/serverdefs
@comment for Xwnmo
lib/wnn/ja_JP/uumrc_vi.xim
lib/wnn/ja_JP/uumrc_e.xim
lib/wnn/ja_JP/uumrc.xim
lib/wnn/ja_JP/uumkey.xim
lib/wnn/ja_JP/uumkey_v.xim
lib/wnn/ja_JP/uumkey_e.xim
lib/wnn/ja_JP/xim.msg
lib/wnn/ja_JP/xjutil.msg
lib/wnn/ko_KR/uumkey.xim
lib/wnn/ko_KR/uumrc.xim
lib/wnn/ko_KR/xim.msg
lib/wnn/ko_KR/xjutil.msg
lib/wnn/lt_LN/uumkey.xim
lib/wnn/lt_LN/uumrc.xim
lib/wnn/lt_LN/xim.msg
lib/wnn/zh_CN/uumkey.xim
lib/wnn/zh_CN/uumkey_p.xim
lib/wnn/zh_CN/uumrc.xim
lib/wnn/zh_CN/uumrc_p.xim
lib/wnn/zh_CN/uumrc_z.xim
lib/wnn/zh_CN/xim.msg
lib/wnn/zh_CN/xjutil.msg
lib/wnn/zh_TW/uumkey.xim
lib/wnn/zh_TW/uumkey_p.xim
lib/wnn/zh_TW/uumrc.xim
lib/wnn/zh_TW/uumrc_p.xim
lib/wnn/zh_TW/uumrc_z.xim
lib/wnn/zh_TW/xim.msg
lib/wnn/zh_TW/xjutil.msg
lib/wnn/cvt_xim_tbl
lib/wnn/ximconf
lib/wnn/ximrc
lib/wnn/ximrc_vi
@comment END for Xwnmo
@dirrm lib/wnn/ko_KR/rk
@dirrm lib/wnn/ko_KR/dic/sys
@dirrm lib/wnn/ko_KR/dic/usr
@dirrm lib/wnn/ko_KR/dic
@dirrm lib/wnn/ko_KR
@comment for Xwnmo
@dirrm lib/wnn/ja_JP
@dirrm lib/wnn/lt_LN
@dirrm lib/wnn/zh_CN
@dirrm lib/wnn/zh_TW
@dirrm lib/wnn
@comment END for Xwnmo
bin/kWnn4/katod
bin/kWnn4/katof
bin/kWnn4/kdtoa
bin/kWnn4/kserver
bin/kWnn4/kuum
bin/kWnn4/kwddel
bin/kWnn4/kwdreg
bin/kWnn4/kwnnkill
bin/kWnn4/kwnnstat
bin/kWnn4/kwnntouch
@dirrm bin/kWnn4
@cwd /usr/X11R6
@comment for Xwnmo
bin/xwnmo
bin/killxwnmo
bin/xjutil
lib/X11/app-defaults/Xwnmo
@comment END for Xwnmo
lib/libkwnn4.so
lib/libkwnn4.so.1.0
@exec /usr/bin/env OBJFORMAT=%%PORTOBJFORMAT%% /sbin/ldconfig -m %D/lib
@unexec /usr/bin/env OBJFORMAT=%%PORTOBJFORMAT%% /sbin/ldconfig -R
lib/libkwnn.a
include/kwnn/commonhd.h
include/kwnn/cplib.h
include/kwnn/jllib.h
include/kwnn/jslib.h
include/kwnn/msg.h
include/kwnn/rk_spclval.h
include/kwnn/wnnerror.h
@dirrm include/kwnn
@comment for Xwnmo
man/man1/xwnmo.1.gz
man/man1/ximrc.1.gz
man/man1/ximconf.1.gz
man/man1/cvt_fun_tbl.1.gz
man/man1/cvt_meta_tbl.1.gz
man/man1/killxwnmo.1.gz
man/man1/cvt_xim_tbl.1.gz
man/man1/xjutil.1.gz
@comment END for Xwnmo
