# net-p2p/btsync/files/openrc-btsync.in
#!/sbin/openrc-run
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# %%PORTNAME%%_user:            The user account %%PORTNAME%% daemon runs as
#                               what you want it to be. It uses '%%USER%%'
#                               user by default. Do not sets it as empty or it
#                               will run as root.
# %%PORTNAME%%_group:           The group account %%PORTNAME%% daemon runs as
#                               what you want it to be. It uses '%%GROUP%%'
#                               group by default. Do not sets it as empty or it
#                               will run as wheel.

name="%%PORTNAME%%"

: ${%%PORTNAME%%_user:="%%USER%%"}
: ${%%PORTNAME%%_group:="%%GROUP%%"}

pidfile="/var/run/%%PORTNAME%%/%%PORTNAME%%.pid"
command="%%PREFIX%%/bin/btsync"
command_args="--config %%PREFIX%%/etc/%%PORTNAME%%.conf"

start_pre() {
	if [ ! -d ${pidfile%/*} ]; then
		install -d -o $%%PORTNAME%%_user -g $%%PORTNAME%%_group ${pidfile%/*}
	fi
}

depend() {
		need localmount
		keyword -shutdown
		after bootmisc
}
