# dns/totd/pkg-plist
sbin/totd
etc/rc.d/totd.sh
etc/totd.conf.sample
