# sysutils/openhpi/pkg-plist
bin/hpi_shell
bin/hpialarms
bin/hpicrypt
bin/hpidomain
bin/hpiel
bin/hpievents
bin/hpifan
bin/hpigensimdata
bin/hpiinv
bin/hpionIBMblade
bin/hpipower
bin/hpireset
bin/hpisensor
bin/hpisettime
bin/hpithres
bin/hpitop
bin/hpitree
bin/hpiwdt
bin/hpixml
bin/ohdomainlist
bin/ohhandler
bin/ohparam
@comment etc/init.d/openhpid
@sample %%ETCDIR%%/openhpi.conf.sample
@sample %%ETCDIR%%/openhpiclient.conf.sample
@sample %%ETCDIR%%/simulation.data.sample
include/openhpi/SaHpi.h
include/openhpi/SaHpiAtca.h
include/openhpi/SaHpiBladeCenter.h
include/openhpi/SaHpiXtca.h
include/openhpi/announcement_utils.h
include/openhpi/el_utils.h
include/openhpi/epath_utils.h
include/openhpi/event_utils.h
include/openhpi/oHpi.h
include/openhpi/oh_clients.h
include/openhpi/oh_utils.h
include/openhpi/rpt_utils.h
include/openhpi/sahpi_enum_utils.h
include/openhpi/sahpi_event_encode.h
include/openhpi/sahpi_event_utils.h
include/openhpi/sahpi_gcrypt_utils.h
include/openhpi/sahpi_struct_utils.h
include/openhpi/sahpi_time_utils.h
include/openhpi/sahpiatca_enum_utils.h
include/openhpi/sahpixtca_enum_utils.h
include/openhpi/uid_utils.h
lib/libopenhpi.a
lib/libopenhpi.so
lib/libopenhpi.so.3
lib/libopenhpi.so.3.6.1
lib/libopenhpi_snmp.a
lib/libopenhpi_snmp.so
lib/libopenhpi_snmp.so.3
lib/libopenhpi_snmp.so.3.6.1
lib/libopenhpi_ssl.a
lib/libopenhpi_ssl.so
lib/libopenhpi_ssl.so.3
lib/libopenhpi_ssl.so.3.6.1
lib/libopenhpimarshal.a
lib/libopenhpimarshal.so
lib/libopenhpimarshal.so.3
lib/libopenhpimarshal.so.3.6.1
lib/libopenhpitransport.a
lib/libopenhpitransport.so
lib/libopenhpitransport.so.3
lib/libopenhpitransport.so.3.6.1
lib/libopenhpiutils.a
lib/libopenhpiutils.so
lib/libopenhpiutils.so.3
lib/libopenhpiutils.so.3.6.1
lib/openhpi/libdyn_simulator.a
lib/openhpi/libdyn_simulator.so
lib/openhpi/libdyn_simulator.so.3
lib/openhpi/libdyn_simulator.so.3.6.1
lib/openhpi/libilo2_ribcl.a
lib/openhpi/libilo2_ribcl.so
lib/openhpi/libilo2_ribcl.so.3
lib/openhpi/libilo2_ribcl.so.3.6.1
lib/openhpi/libipmidirect.a
lib/openhpi/libipmidirect.so
lib/openhpi/libipmidirect.so.3
lib/openhpi/libipmidirect.so.3.6.1
lib/openhpi/liboa_soap.a
lib/openhpi/liboa_soap.so
lib/openhpi/liboa_soap.so.3
lib/openhpi/liboa_soap.so.3.6.1
lib/openhpi/libsimulator.a
lib/openhpi/libsimulator.so
lib/openhpi/libsimulator.so.3
lib/openhpi/libsimulator.so.3.6.1
lib/openhpi/libslave.a
lib/openhpi/libslave.so
lib/openhpi/libslave.so.3
lib/openhpi/libslave.so.3.6.1
lib/openhpi/libsnmp_bc.a
lib/openhpi/libsnmp_bc.so
lib/openhpi/libsnmp_bc.so.3
lib/openhpi/libsnmp_bc.so.3.6.1
lib/openhpi/libtest_agent.a
lib/openhpi/libtest_agent.so
lib/openhpi/libtest_agent.so.3
lib/openhpi/libtest_agent.so.3.6.1
libdata/pkgconfig/openhpi.pc
libdata/pkgconfig/openhpiutils.pc
man/man1/hpi_shell.1.gz
man/man1/hpialarms.1.gz
man/man1/hpicrypt.1.gz
man/man1/hpidomain.1.gz
man/man1/hpiel.1.gz
man/man1/hpievents.1.gz
man/man1/hpifan.1.gz
man/man1/hpigensimdata.1.gz
man/man1/hpiinv.1.gz
man/man1/hpionIBMblade.1.gz
man/man1/hpipower.1.gz
man/man1/hpireset.1.gz
man/man1/hpisensor.1.gz
man/man1/hpisettime.1.gz
man/man1/hpithres.1.gz
man/man1/hpitop.1.gz
man/man1/hpitree.1.gz
man/man1/hpiwdt.1.gz
man/man1/hpixml.1.gz
man/man1/ohdomainlist.1.gz
man/man1/ohhandler.1.gz
man/man1/ohparam.1.gz
man/man7/openhpi.7.gz
man/man8/openhpid.8.gz
sbin/openhpid
@comment %%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/ChangeLog
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/README.daemon
@dir /var/lib/openhpi
@dir /var/lib
