# sysutils/nut22/pkg-plist
bin/upsc
bin/upscmd
bin/upslog
bin/upsrw
etc/nut/ups.conf.sample
etc/nut/upsd.conf.sample
etc/nut/upsd.users.sample
etc/nut/upsmon.conf.sample
etc/nut/upssched.conf.sample
etc/rc.d/nut.sh.sample
libexec/nut/apcsmart
libexec/nut/bcmxcp
libexec/nut/belkin
libexec/nut/bestfortress
libexec/nut/bestuferrups
libexec/nut/bestups
libexec/nut/blazer
libexec/nut/cyberpower
libexec/nut/esupssmart
libexec/nut/etapro
libexec/nut/everups
libexec/nut/fentonups
libexec/nut/ferrups
libexec/nut/genericups
libexec/nut/hp
libexec/nut/isbmex
libexec/nut/liebert
libexec/nut/masterguard
libexec/nut/mge-shut
libexec/nut/mge-utalk
libexec/nut/microdowell
libexec/nut/newapc
libexec/nut/oneac
libexec/nut/powercom
libexec/nut/safenet
libexec/nut/sec
libexec/nut/sms
libexec/nut/tripplite
libexec/nut/tripplitesu
libexec/nut/upsdrvctl
libexec/nut/victronups
sbin/upsd
sbin/upsmon
sbin/upssched
share/cmdvartab
%%CGIFILES%%upsimage.cgi
%%CGIFILES%%upsset.cgi
%%CGIFILES%%upsstats.cgi
%%CGIETC%%hosts.conf.sample
%%CGIETC%%upsset.conf.sample
%%CGIETC%%upsstats.html.sample
%%CGIETC%%upsstats-single.html.sample
%%PORTDOCS%%%%DOCSDIR%%/CHANGES
%%PORTDOCS%%%%DOCSDIR%%/CREDITS
%%PORTDOCS%%%%DOCSDIR%%/FAQ
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/NEWS
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/UPGRADING
%%PORTDOCS%%%%DOCSDIR%%/big-servers.txt
%%PORTDOCS%%%%DOCSDIR%%/cables/apc.txt
%%PORTDOCS%%%%DOCSDIR%%/cables/ge-imv-victron.txt
%%PORTDOCS%%%%DOCSDIR%%/cables/imv.txt
%%PORTDOCS%%%%DOCSDIR%%/cables/mgeups.txt
%%PORTDOCS%%%%DOCSDIR%%/cables/powerware.txt
%%PORTDOCS%%%%DOCSDIR%%/cables/repotec.txt
%%PORTDOCS%%%%DOCSDIR%%/cables/sms.txt
%%PORTDOCS%%%%DOCSDIR%%/chroot.txt
%%PORTDOCS%%%%DOCSDIR%%/commands.txt
%%PORTDOCS%%%%DOCSDIR%%/config-files.txt
%%PORTDOCS%%%%DOCSDIR%%/configure.txt
%%PORTDOCS%%%%DOCSDIR%%/contact-closure.txt
%%PORTDOCS%%%%DOCSDIR%%/data-room.txt
%%PORTDOCS%%%%DOCSDIR%%/design.txt
%%PORTDOCS%%%%DOCSDIR%%/developers.txt
%%PORTDOCS%%%%DOCSDIR%%/driver.list
%%PORTDOCS%%%%DOCSDIR%%/new-drivers.txt
%%PORTDOCS%%%%DOCSDIR%%/new-names.txt
%%PORTDOCS%%%%DOCSDIR%%/osd-notify.txt
%%PORTDOCS%%%%DOCSDIR%%/pager.txt
%%PORTDOCS%%%%DOCSDIR%%/powersaving.txt
%%PORTDOCS%%%%DOCSDIR%%/protocol.txt
%%PORTDOCS%%%%DOCSDIR%%/shutdown.txt
%%PORTDOCS%%%%DOCSDIR%%/snmp.txt
%%PORTDOCS%%%%DOCSDIR%%/sock-protocol.txt
%%PORTDOCS%%%%DOCSDIR%%/ssl.txt
%%PORTDOCS%%%%DOCSDIR%%/todo.txt
%%PORTDOCS%%%%DOCSDIR%%/upssched.txt
@dirrm %%DOCSDIR%%/cables
@dirrm %%DOCSDIR%%/drivers
@dirrm %%DOCSDIR%%
@dirrm libexec/nut
%%CGIDIRRM%%
@unexec rmdir %D/etc/nut 2>/dev/null || true
