# net-mgmt/zenoss/files/zenoss.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: zenoss
# REQUIRE: login
#
# Add the following line to /etc/rc.conf to enable zenoss:
#
# zenoss_enable="YES"

command="%%ZENHOME%%/bin/zenoss"

. /etc/rc.subr

load_rc_config zenoss

#set defaults
zenoss_enable=${zenoss_enable:-"NO"}

name=zenoss
rcvar=`set_rcvar`
extra_commands="init"

start_cmd="zenoss_command start"
stop_cmd="zenoss_command stop"
restart_cmd="zenoss_command restart"
status_cmd="zenoss_command status"

init_cmd="zenoss_init"

zenoss_command()
{
   su -l zenoss -c "exec ${command} ${rc_arg}"
}

zenoss_init()
{
   su -l zenoss -c "cd %%ZENHOME%%/init && ./init.sh"
}

run_rc_command "$1"
