# audio/shairport/files/openrc-shairport.in
#!/sbin/openrc-run
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf to
# enable shairport:
#
# shairport_flags (flags):	Set extra flags here.
#				shairport_flags="-a name" to change Airport name
#				Default is empty "".
# shairport_user (user):	Set user to run shairport.
#				Default is "nobody".

name="shairport"

: ${shairport_user="nobody"}

command="%%PREFIX%%/bin/${name}"
pidfile="/var/run/${name}/${name}.pid"

command_args="-d -P ${pidfile}"
command_user=${shairport_user}

start_pre() {
	install -d -o ${shairport_user} -g wheel -m 755 /var/run/${name}
}

depend() {
	provide shairport
	#need avahi_daemon
	keyword -shutdown
}
