# audio/shairport/files/shairport.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: shairport
# REQUIRE: LOGIN avahi_daemon
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf to
# enable shairport:
#
# shairport_enable (bool):	Set it to "YES" to enable shaiport.
#				Default is "NO".
# shairport_flags (flags):	Set extra flags here.
#				shairport_flags="-a name" to change Airport name
#				Default is empty "".
# shairport_user (user):	Set user to run shairport.
#				Default is "nobody".

. /etc/rc.subr

name="shairport"
rcvar=shairport_enable

load_rc_config ${name}

: ${shairport_enable="NO"}
: ${shairport_user="nobody"}

command="%%PREFIX%%/bin/${name}"
pidfile="/var/run/${name}/${name}.pid"

command_args="-d -P ${pidfile}"
start_precmd="install -d -o ${shairport_user} -g wheel -m 755 /var/run/${name}"

run_rc_command "$1"
