# net-im/jggtrans/files/jggtrans.sh.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: jggtrans
# REQUIRE: DAEMON %%JABBER_USER%%
# BEFORE: LOGIN
# KEYWORD: shutdown

# DO NOT CHANGE THESE DEFAULT VALUES HERE
# SET THEM IN THE /etc/rc.conf 

. /etc/rc.subr

name="jggtrans"
rcvar=jggtrans_enable

load_rc_config $name

: ${jggtrans_enable="NO"}
: ${jggtrans_program="%%PREFIX%%/sbin/jggtrans"}
: ${jggtrans_flags="-u %%JABBER_USER%% -g %%JABBER_USER%%"}
: ${jggtrans_pidfile="/var/jabberd/pid/jggtrans.pid"}

start_cmd="echo \"Starting ${name}.\"; /usr/bin/nice -5 ${jggtrans_program} ${jggtrans_flags}"

run_rc_command "$1"
