# security/vaultwarden/pkg-plist
bin/vaultwarden
etc/rc.conf.d/vaultwarden.sample
etc/rc.d/vaultwarden
@dir %%ETCDIR%%/rc.conf.d
@owner www
@group www
@mode 755
@dir %%WWWDIR%%/data
