# sysutils/iocage/pkg-plist
etc/rc.d/iocage
lib/iocage/ioc-cmd
lib/iocage/ioc-common
lib/iocage/ioc-configure
lib/iocage/ioc-globals
lib/iocage/ioc-hacks
lib/iocage/ioc-help
lib/iocage/ioc-image
lib/iocage/ioc-info
lib/iocage/ioc-network
lib/iocage/ioc-rc
lib/iocage/ioc-rctl
lib/iocage/ioc-snapshot
lib/iocage/ioc-upgrade
lib/iocage/ioc-zfs
man/man8/iocage.8.gz
sbin/iocage
