# security/greenbone-security-assistant/files/openrc-gsad.in
#!/sbin/openrc-run

name=gsad
command="%%PREFIX%%/sbin/gsad"
pidfile=/var/run/gsad.pid

depend() {
		keyword -shutdown
}
