# net-mgmt/netxms/files/openrc-netxmsd.in
#!/sbin/openrc-run

name="netxmsd"

command="%%PREFIX%%/bin/netxmsd"
required_files="%%PREFIX%%/etc/netxmsd.conf"
command_args="-c %%PREFIX%%/etc/netxmsd.conf -d -e"

depend() {
	keyword -shutdown
}
