# sysutils/fcron/files/openrc-fcron.in
#!/sbin/openrc-run

description="Periodic command scheduler"
name=fcron

command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/${name}.pid"

depend() {
	keyword -shutdown
}
