# sysutils/bsdstats/files/openrc-bsdstats.in
#!/sbin/openrc-run

name=bsdstats
command="%%PREFIX%%/etc/periodic/monthly/300.statistics"
command_args="-- -shortdelay >/dev/null 2>/dev/null &"

depend() {
	need localmount
	use net
	keyword -stop -shutdown
}
