# sysutils/bsdstats/files/openrc-bsdstats.in
#!/sbin/openrc-run

name=bsdstats

depend()
{
	need localmount
	use net
	keyword -stop -shutdown
}

start()
{
  ebegin "Starting $name"
  (%%PREFIX%%/etc/periodic/monthly/300.statistics -- -shortdelay >/dev/null 2>/dev/null) &
  mark_service_started
  eend $?
}

stop()
{
  return 0
}
