# japanese/vje30/pkg-plist
@unexec for f in %%BINSTUFF%%; do rm %D/bin/$f; done
etc/rc.d/vje.sh
etc/vje30/vje.cfg
vje30/acc/default.kig
vje30/acc/vjed95f1.dat
vje30/acc/vjed95f2.dat
vje30/acc/vjed95f3.dat
vje30/acc/vjed95f4.dat
vje30/acc/vjed95f5.dat
vje30/acc/vjed95f6.dat
vje30/acc/vjed95f7.dat
vje30/acc/vjed95f8.dat
vje30/acc/vjed95f9.dat
vje30/acc/vjed95fb.dat
vje30/acc/vjed95fg.dat
vje30/acc/vjed95fk.dat
vje30/acc/vjed95ft.dat
vje30/acc/vjed95fx.dat
vje30/dic/vjed95j.dic
vje30/dic/vjed95e.dat
vje30/dic/vjed95k.dic
vje30/dic/vjed95m.dic
vje30/dic/vjed95n.dic
vje30/dic/vjed95t.dic
vje30/dic/vjed95ye.dic
vje30/dic/vjed95yw.dic
vje30/dic/sample.euc
vje30/com/general.sty
vje30/com/vjed95u.dic
%%PORTDOCS%%vje30/doc/delta.htm
%%PORTDOCS%%vje30/doc/images/r4_7_4.gif
%%PORTDOCS%%vje30/doc/images/bc01.gif
%%PORTDOCS%%vje30/doc/images/bc02.gif
%%PORTDOCS%%vje30/doc/images/bc03.gif
%%PORTDOCS%%vje30/doc/images/bc04.gif
%%PORTDOCS%%vje30/doc/images/bc05.gif
%%PORTDOCS%%vje30/doc/images/bc06.gif
%%PORTDOCS%%vje30/doc/images/bc07.gif
%%PORTDOCS%%vje30/doc/images/bc08.gif
%%PORTDOCS%%vje30/doc/images/bc10.gif
%%PORTDOCS%%vje30/doc/images/bc11.gif
%%PORTDOCS%%vje30/doc/images/bc12.gif
%%PORTDOCS%%vje30/doc/images/bc13.gif
%%PORTDOCS%%vje30/doc/images/bc14.gif
%%PORTDOCS%%vje30/doc/images/bc15.gif
%%PORTDOCS%%vje30/doc/images/bc16.gif
%%PORTDOCS%%vje30/doc/images/bc17.gif
%%PORTDOCS%%vje30/doc/images/bc17a.gif
%%PORTDOCS%%vje30/doc/images/bc18.gif
%%PORTDOCS%%vje30/doc/images/bc19.gif
%%PORTDOCS%%vje30/doc/images/bc19a.gif
%%PORTDOCS%%vje30/doc/images/bc20.gif
%%PORTDOCS%%vje30/doc/images/bc21.gif
%%PORTDOCS%%vje30/doc/images/bc21a.gif
%%PORTDOCS%%vje30/doc/images/bc22.gif
%%PORTDOCS%%vje30/doc/images/bc22a.gif
%%PORTDOCS%%vje30/doc/images/bc23.gif
%%PORTDOCS%%vje30/doc/images/bc24.gif
%%PORTDOCS%%vje30/doc/images/bc25.gif
%%PORTDOCS%%vje30/doc/images/bc26.gif
%%PORTDOCS%%vje30/doc/images/bc27.gif
%%PORTDOCS%%vje30/doc/images/bc28.gif
%%PORTDOCS%%vje30/doc/images/bc61.gif
%%PORTDOCS%%vje30/doc/images/bc62.gif
%%PORTDOCS%%vje30/doc/images/bc63.gif
%%PORTDOCS%%vje30/doc/images/bc64.gif
%%PORTDOCS%%vje30/doc/images/bc65.gif
%%PORTDOCS%%vje30/doc/images/bc66.gif
%%PORTDOCS%%vje30/doc/images/bc67.gif
%%PORTDOCS%%vje30/doc/images/bceuc.gif
%%PORTDOCS%%vje30/doc/images/bcjis.gif
%%PORTDOCS%%vje30/doc/images/bckuten.gif
%%PORTDOCS%%vje30/doc/images/bcsjis.gif
%%PORTDOCS%%vje30/doc/images/delta2ux.gif
%%PORTDOCS%%vje30/doc/images/g1_1.gif
%%PORTDOCS%%vje30/doc/images/g2_1_1.gif
%%PORTDOCS%%vje30/doc/images/g2_3.gif
%%PORTDOCS%%vje30/doc/images/g2_4.gif
%%PORTDOCS%%vje30/doc/images/g2_5.gif
%%PORTDOCS%%vje30/doc/images/r1_1_1.gif
%%PORTDOCS%%vje30/doc/images/r1_1_2.gif
%%PORTDOCS%%vje30/doc/images/r1_1_3.gif
%%PORTDOCS%%vje30/doc/images/r1_2_1.gif
%%PORTDOCS%%vje30/doc/images/r1_2_2.gif
%%PORTDOCS%%vje30/doc/images/r1_2_3.gif
%%PORTDOCS%%vje30/doc/images/r1_2_4.gif
%%PORTDOCS%%vje30/doc/images/r1_3_1.gif
%%PORTDOCS%%vje30/doc/images/r1_3_2.gif
%%PORTDOCS%%vje30/doc/images/r2_3_1.gif
%%PORTDOCS%%vje30/doc/images/r2_4_1.gif
%%PORTDOCS%%vje30/doc/images/r4_1_1.gif
%%PORTDOCS%%vje30/doc/images/r4_1_2.gif
%%PORTDOCS%%vje30/doc/images/r4_1_3.gif
%%PORTDOCS%%vje30/doc/images/r4_1_4.gif
%%PORTDOCS%%vje30/doc/images/r4_10_1.gif
%%PORTDOCS%%vje30/doc/images/r4_10_2.gif
%%PORTDOCS%%vje30/doc/images/r4_10_3.gif
%%PORTDOCS%%vje30/doc/images/r4_11_1.gif
%%PORTDOCS%%vje30/doc/images/r4_11_2.gif
%%PORTDOCS%%vje30/doc/images/r4_12_1.gif
%%PORTDOCS%%vje30/doc/images/r4_12_2.gif
%%PORTDOCS%%vje30/doc/images/r4_12_3.gif
%%PORTDOCS%%vje30/doc/images/r4_13_1.gif
%%PORTDOCS%%vje30/doc/images/r4_13_2.gif
%%PORTDOCS%%vje30/doc/images/r4_2_1.gif
%%PORTDOCS%%vje30/doc/images/r4_2_2.gif
%%PORTDOCS%%vje30/doc/images/r4_3_1.gif
%%PORTDOCS%%vje30/doc/images/r4_3_2.gif
%%PORTDOCS%%vje30/doc/images/r4_4_1.gif
%%PORTDOCS%%vje30/doc/images/r4_4_2.gif
%%PORTDOCS%%vje30/doc/images/r4_5_1.gif
%%PORTDOCS%%vje30/doc/images/r4_5_2.gif
%%PORTDOCS%%vje30/doc/images/r4_5_3.gif
%%PORTDOCS%%vje30/doc/images/r4_5_4.gif
%%PORTDOCS%%vje30/doc/images/r4_6_1.gif
%%PORTDOCS%%vje30/doc/images/r4_6_2.gif
%%PORTDOCS%%vje30/doc/images/r4_7_1.gif
%%PORTDOCS%%vje30/doc/images/r4_7_2.gif
%%PORTDOCS%%vje30/doc/images/r4_7_3.gif
%%PORTDOCS%%vje30/doc/images/batu.gif
%%PORTDOCS%%vje30/doc/images/r4_7_5.gif
%%PORTDOCS%%vje30/doc/images/r4_7_6.gif
%%PORTDOCS%%vje30/doc/images/r4_8_1.gif
%%PORTDOCS%%vje30/doc/images/r4_8_2.gif
%%PORTDOCS%%vje30/doc/images/r4_8_3.gif
%%PORTDOCS%%vje30/doc/images/r4_9_1.gif
%%PORTDOCS%%vje30/doc/images/r4_9_2.gif
%%PORTDOCS%%vje30/doc/images/r4_9_3.gif
%%PORTDOCS%%vje30/doc/images/toolbar1.gif
%%PORTDOCS%%vje30/doc/images/bc29.gif
%%PORTDOCS%%vje30/doc/images/toolbar2.gif
%%PORTDOCS%%vje30/doc/images/toolbar3.gif
%%PORTDOCS%%vje30/doc/images/vacslogo.gif
%%PORTDOCS%%vje30/doc/images/bc30.gif
%%PORTDOCS%%vje30/doc/images/bc31.gif
%%PORTDOCS%%vje30/doc/images/bc32.gif
%%PORTDOCS%%vje30/doc/images/bc33.gif
%%PORTDOCS%%vje30/doc/images/bc34.gif
%%PORTDOCS%%vje30/doc/images/bc35.gif
%%PORTDOCS%%vje30/doc/images/bc36.gif
%%PORTDOCS%%vje30/doc/images/bc37.gif
%%PORTDOCS%%vje30/doc/images/bc38.gif
%%PORTDOCS%%vje30/doc/images/bc39.gif
%%PORTDOCS%%vje30/doc/images/bc40.gif
%%PORTDOCS%%vje30/doc/images/bc41.gif
%%PORTDOCS%%vje30/doc/images/bc42.gif
%%PORTDOCS%%vje30/doc/images/bc43.gif
%%PORTDOCS%%vje30/doc/images/bc44.gif
%%PORTDOCS%%vje30/doc/images/bc45.gif
%%PORTDOCS%%vje30/doc/images/bc46.gif
%%PORTDOCS%%vje30/doc/images/bc47.gif
%%PORTDOCS%%vje30/doc/images/bc48.gif
%%PORTDOCS%%vje30/doc/images/bc49.gif
%%PORTDOCS%%vje30/doc/images/bc50.gif
%%PORTDOCS%%vje30/doc/images/bc51.gif
%%PORTDOCS%%vje30/doc/images/bc52.gif
%%PORTDOCS%%vje30/doc/images/bc53.gif
%%PORTDOCS%%vje30/doc/images/bc54.gif
%%PORTDOCS%%vje30/doc/images/bc55.gif
%%PORTDOCS%%vje30/doc/images/bc56.gif
%%PORTDOCS%%vje30/doc/images/bc57.gif
%%PORTDOCS%%vje30/doc/images/bc58.gif
%%PORTDOCS%%vje30/doc/images/bc59.gif
%%PORTDOCS%%vje30/doc/images/bc60.gif
%%PORTDOCS%%vje30/doc/images/mark01.gif
%%PORTDOCS%%vje30/doc/images/mark02.gif
%%PORTDOCS%%vje30/doc/images/r3_2_7.gif
%%PORTDOCS%%vje30/doc/images/r3_3_1.gif
%%PORTDOCS%%vje30/doc/images/r3_2_1.gif
%%PORTDOCS%%vje30/doc/images/r3_4_1.gif
%%PORTDOCS%%vje30/doc/images/r3_3_2.gif
%%PORTDOCS%%vje30/doc/images/r3_2_2.gif
%%PORTDOCS%%vje30/doc/images/r3_2_3.gif
%%PORTDOCS%%vje30/doc/images/r3_2_4.gif
%%PORTDOCS%%vje30/doc/images/r3_2_5.gif
%%PORTDOCS%%vje30/doc/images/r3_2_6.gif
%%PORTDOCS%%vje30/doc/images/r3_2_8.gif
%%PORTDOCS%%vje30/doc/images/r3_3_3.gif
%%PORTDOCS%%vje30/doc/images/r3_3_4.gif
%%PORTDOCS%%vje30/doc/images/r4_14_3.gif
%%PORTDOCS%%vje30/doc/images/r66_1_2.gif
%%PORTDOCS%%vje30/doc/images/r4_14_4.gif
%%PORTDOCS%%vje30/doc/images/r4_14_2.gif
%%PORTDOCS%%vje30/doc/images/r66_1_1.gif
%%PORTDOCS%%vje30/doc/images/r4_14_1.gif
%%PORTDOCS%%vje30/doc/images/r4_14_5.gif
%%PORTDOCS%%vje30/doc/images/sb01.gif
%%PORTDOCS%%vje30/doc/images/sb02.gif
%%PORTDOCS%%vje30/doc/images/sb03.gif
%%PORTDOCS%%@dirrm vje30/doc/images
%%PORTDOCS%%@dirrm vje30/doc
vje30/env/default.ach
vje30/env/default.kch
vje30/env/deltasr.key
vje30/env/vje.clr
vje30/env/aiueo.kch
vje30/env/canna.key
vje30/env/vje.rom
vje30/env/abc.ach
vje30/env/ascii.ach
vje30/env/dvorak.ach
vje30/env/jis.ach
vje30/env/atok.rom
vje30/env/msime.rom
vje30/env/wnn.rom
vje30/env/atok.key
vje30/env/msime.key
vje30/env/wxg.key
vje30/env/wxg.rom
vje30/env/canna.rom
vje30/env/wnn.key
vje30/lib/libvjed30f.so
vje30/lib/libvjed30k.so
%%PORTDOCS%%vje30/README
%%PORTDOCS%%vje30/FAQ
%%PORTDOCS%%vje30/REQUEST
vje30/.version
vje30/bin/vje
vje30/bin/vjed
vje30/bin/vjekill
vje30/bin/vadd
vje30/bin/vdel
vje30/bin/vmaked
vje30/bin/vdispd
vje30/bin/vpen
vje30/bin/vpr
vje30/bin/vprc
vje30/bin/vpu
vje30/bin/vjeacc
vje30/bin/vjeacc.jar
@dirrm etc/vje30
@dirrm vje30/acc
@dirrm vje30/bin
@dirrm vje30/com
@dirrm vje30/dic
@dirrm vje30/env
@dirrm vje30/lib
@dirrm vje30
@exec for f in %%BINSTUFF%%; do ln -fs %D/vje30/bin/$f %D/bin/$f; done
@exec chmod a+w %D/vje30/dic/vjed95m.dic
@exec chmod a+x %D/etc/rc.d/vje.sh
@unexec cp -Pp /etc/services /etc/services.bak ; grep -vw '^vjed' /etc/services.bak > /etc/services
@exec cp -Pp /etc/services /etc/services.bak ; grep -qw '^vjed' /etc/services || echo "vjed		11493/tcp	# VJE-Delta Server" >> /etc/services
@unexec [ -L /usr/local/vje30 ] && rm /usr/local/vje30 || true
@unexec [ -L /usr/local/etc/vje30 ] && rm /usr/local/etc/vje30 || true
@exec [ X"%D" != X"/usr/local" ] && (ln -s %D/vje30 /usr/local/; ln -s %D/etc/vje30 /usr/local/etc/) || true
