# net/openldap/pkg-plist
bin/ldapadd
bin/ldapdelete
bin/ldapmodify
bin/ldapmodrdn
bin/ldappasswd
bin/ldapsearch
bin/ud
@unexec if cmp -s %D/etc/openldap/ldap.conf %D/etc/openldap/ldap.conf.default; then rm -f %D/etc/openldap/ldap.conf; fi
etc/openldap/ldap.conf.default
@exec [ ! -f %B/ldap.conf ] && cp %B/%f %B/ldap.conf
@unexec if cmp -s %D/etc/openldap/ldapfilter.conf %D/etc/openldap/ldapfilter.conf.default; then rm -f %D/etc/openldap/ldapfilter.conf; fi
etc/openldap/ldapfilter.conf.default
@exec [ ! -f %B/ldapfilter.conf ] && cp %B/%f %B/ldapfilter.conf
@unexec if cmp -s %D/etc/openldap/ldapsearchprefs.conf %D/etc/openldap/ldapsearchprefs.conf.default; then rm -f %D/etc/openldap/ldapsearchprefs.conf; fi
etc/openldap/ldapsearchprefs.conf.default
@exec [ ! -f %B/ldapsearchprefs.conf ] && cp %B/%f %B/ldapsearchprefs.conf
@unexec if cmp -s %D/etc/openldap/ldaptemplates.conf %D/etc/openldap/ldaptemplates.conf.default; then rm -f %D/etc/openldap/ldaptemplates.conf; fi
etc/openldap/ldaptemplates.conf.default
@exec [ ! -f %B/ldaptemplates.conf ] && cp %B/%f %B/ldaptemplates.conf
@unexec if cmp -s %D/etc/openldap/slapd.at.conf %D/etc/openldap/slapd.at.conf.default; then rm -f %D/etc/openldap/slapd.at.conf; fi
etc/openldap/slapd.at.conf.default
@exec [ ! -f %B/slapd.at.conf ] && cp %B/%f %B/slapd.at.conf
@unexec if cmp -s %D/etc/openldap/slapd.conf %D/etc/openldap/slapd.conf.default; then rm -f %D/etc/openldap/slapd.conf; fi
etc/openldap/slapd.conf.default
@exec [ ! -f %B/slapd.conf ] && cp %B/%f %B/slapd.conf
@unexec if cmp -s %D/etc/openldap/slapd.oc.conf %D/etc/openldap/slapd.oc.conf.default; then rm -f %D/etc/openldap/slapd.oc.conf; fi
etc/openldap/slapd.oc.conf.default
@exec [ ! -f %B/slapd.oc.conf ] && cp %B/%f %B/slapd.oc.conf
@dirrm etc/openldap
etc/rc.d/slapd.sh.sample
include/disptmpl.h
include/lber.h
include/ldap.h
include/ldap_cdefs.h
include/srchpref.h
lib/liblber.a
lib/liblber.so
lib/liblber.so.1
lib/libldap.a
lib/libldap.so
lib/libldap.so.1
libexec/fax500
libexec/go500
libexec/go500gw
libexec/in.xfingerd
libexec/mail500
libexec/rcpt500
libexec/rp500
libexec/slapd
libexec/slurpd
libexec/xrpcomp
sbin/centipede
sbin/ldbmcat
sbin/ldbmtest
sbin/ldif
sbin/ldif2id2children
sbin/ldif2id2entry
sbin/ldif2index
sbin/ldif2ldbm
share/openldap/go500gw.help
share/openldap/ldapfriendly
share/openldap/rcpt500.help
@dirrm share/openldap
