# net-im/chattahoochie/pkg-plist
bin/chattahoochie
etc/rc.d/chattahoochie
