# net-im/chattahoochie/pkg-plist
bin/chattahoochie
etc/rc.d/chattahoochie
man/man8/chattahoochie.8.gz
