# net/gnatsd/files/gnatsd.sh.in
#!/bin/sh
#
# PROVIDE: gnatsd
# REQUIRE: DAEMON
#
# Add the following lines to /etc/rc.conf to run NATS:
#
# gnatsd_enable (bool):		Set it to "YES" to enable gnatsd server.
#				Default is "NO".
# gnatsd_user:			User name to run as. default "nobody"
# gnatsd_options:		Options to pass gnatsd server
#

. /etc/rc.subr

name="gnatsd"
rcvar=`set_rcvar`

load_rc_config ${name}
: ${gnatsd_enable="NO"}
: ${gnatsd_user="nobody"}
: ${gnatsd_options="-c %%PREFIX%%/etc/gnatsd.conf"}

command=/usr/sbin/daemon
procname=%%PREFIX%%/bin/gnatsd
command_args="${procname} ${gnatsd_options}"

run_rc_command "$1"
