# misc/yaunc/files/yaunc.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: yaunc
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# yaunc_enable (bool):   Set to NO by default.
#               Set it to YES to enable yaunc.
# yaunc_config (path):   Set to %%PREFIX%%/etc/yaunc.conf
#               by default.

. /etc/rc.subr

name="yaunc"
rcvar=yaunc_enable

command=%%PREFIX%%/bin/${name}

load_rc_config $name

: ${yaunc_enable="NO"}
: ${yaunc_config="%%PREFIX%%/etc/yaunc.conf"}

command_args="-f $yaunc_config"

run_rc_command "$1"
