# databases/memcachedb/files/memcachedb.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: memcachedb
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# memcachedb_enable (bool):	Set to NO by default.
#				Set it to YES to enable memcachedb.
#
. /etc/rc.subr

name="memcachedb"
rcvar=`set_rcvar`

load_rc_config $name

command="%%PREFIX%%/bin/${name}"

: ${memcachedb_flags="-p11211 -l 127.0.0.1 -d -u nobody -f /var/db/memcachedb/default.db -H /var/db/memcachedb -P /var/run/memcachedb.pid"}

run_rc_command "$1"

