# sysutils/bastille/pkg-plist
@sample %%ETCDIR%%/bastille.conf.sample
bin/bastille
etc/rc.d/bastille
%%DATADIR%%/bootstrap.sh
%%DATADIR%%/clone.sh
%%DATADIR%%/cmd.sh
%%DATADIR%%/colors.pre.sh
%%DATADIR%%/common.sh
%%DATADIR%%/config.sh
%%DATADIR%%/console.sh
%%DATADIR%%/convert.sh
%%DATADIR%%/cp.sh
%%DATADIR%%/create.sh
%%DATADIR%%/destroy.sh
%%DATADIR%%/edit.sh
%%DATADIR%%/etcupdate.sh
%%DATADIR%%/export.sh
%%DATADIR%%/htop.sh
%%DATADIR%%/import.sh
%%DATADIR%%/jcp.sh
%%DATADIR%%/limits.sh
%%DATADIR%%/list.sh
%%DATADIR%%/migrate.sh
%%DATADIR%%/monitor.sh
%%DATADIR%%/mount.sh
%%DATADIR%%/network.sh
%%DATADIR%%/pkg.sh
%%DATADIR%%/rcp.sh
%%DATADIR%%/rdr.sh
%%DATADIR%%/rename.sh
%%DATADIR%%/restart.sh
%%DATADIR%%/service.sh
%%DATADIR%%/setup.sh
%%DATADIR%%/start.sh
%%DATADIR%%/stop.sh
%%DATADIR%%/sysrc.sh
%%DATADIR%%/tags.sh
%%DATADIR%%/template.sh
%%DATADIR%%/top.sh
%%DATADIR%%/umount.sh
%%DATADIR%%/update.sh
%%DATADIR%%/upgrade.sh
%%DATADIR%%/verify.sh
%%DATADIR%%/zfs.sh
%%DATADIR%%/pkgbase/FreeBSD-base.conf
%%DATADIR%%/templates/default/base/Bastillefile
%%DATADIR%%/templates/default/clone/Bastillefile
%%DATADIR%%/templates/default/empty/Bastillefile
%%DATADIR%%/templates/default/linux/Bastillefile
%%DATADIR%%/templates/default/thick/Bastillefile
%%DATADIR%%/templates/default/thin/Bastillefile
%%DATADIR%%/templates/default/vlan/Bastillefile
%%DATADIR%%/templates/default/vnet/Bastillefile
share/man/man8/bastille.8.gz
share/man/man8/bastille-bootstrap.8.gz
share/man/man8/bastille-clone.8.gz
share/man/man8/bastille-cmd.8.gz
share/man/man8/bastille-config.8.gz
share/man/man8/bastille-console.8.gz
share/man/man8/bastille-convert.8.gz
share/man/man8/bastille-cp.8.gz
share/man/man8/bastille-create.8.gz
share/man/man8/bastille-destroy.8.gz
share/man/man8/bastille-edit.8.gz
share/man/man8/bastille-etcupdate.8.gz
share/man/man8/bastille-export.8.gz
share/man/man8/bastille-htop.8.gz
share/man/man8/bastille-import.8.gz
share/man/man8/bastille-jcp.8.gz
share/man/man8/bastille-limits.8.gz
share/man/man8/bastille-list.8.gz
share/man/man8/bastille-migrate.8.gz
share/man/man8/bastille-monitor.8.gz
share/man/man8/bastille-mount.8.gz
share/man/man8/bastille-network.8.gz
share/man/man8/bastille-pkg.8.gz
share/man/man8/bastille-rcp.8.gz
share/man/man8/bastille-rdr.8.gz
share/man/man8/bastille-rename.8.gz
share/man/man8/bastille-restart.8.gz
share/man/man8/bastille-service.8.gz
share/man/man8/bastille-setup.8.gz
share/man/man8/bastille-start.8.gz
share/man/man8/bastille-stop.8.gz
share/man/man8/bastille-sysrc.8.gz
share/man/man8/bastille-tags.8.gz
share/man/man8/bastille-template.8.gz
share/man/man8/bastille-top.8.gz
share/man/man8/bastille-umount.8.gz
share/man/man8/bastille-update.8.gz
share/man/man8/bastille-upgrade.8.gz
share/man/man8/bastille-verify.8.gz
share/man/man8/bastille-zfs.8.gz
share/man/man5/bastille.conf.5.gz
@dir %%ETCDIR%%
