# sysutils/bastille/pkg-plist
@sample %%ETCDIR%%/bastille.conf.sample
bin/bastille
etc/rc.d/bastille
%%DATADIR%%/bootstrap.sh
%%DATADIR%%/clone.sh
%%DATADIR%%/cmd.sh
%%DATADIR%%/colors.pre.sh
%%DATADIR%%/common.sh
%%DATADIR%%/config.sh
%%DATADIR%%/console.sh
%%DATADIR%%/convert.sh
%%DATADIR%%/cp.sh
%%DATADIR%%/create.sh
%%DATADIR%%/destroy.sh
%%DATADIR%%/edit.sh
%%DATADIR%%/etcupdate.sh
%%DATADIR%%/export.sh
%%DATADIR%%/htop.sh
%%DATADIR%%/import.sh
%%DATADIR%%/jcp.sh
%%DATADIR%%/limits.sh
%%DATADIR%%/list.sh
%%DATADIR%%/migrate.sh
%%DATADIR%%/mount.sh
%%DATADIR%%/network.sh
%%DATADIR%%/pkg.sh
%%DATADIR%%/rcp.sh
%%DATADIR%%/rdr.sh
%%DATADIR%%/rename.sh
%%DATADIR%%/restart.sh
%%DATADIR%%/service.sh
%%DATADIR%%/setup.sh
%%DATADIR%%/start.sh
%%DATADIR%%/stop.sh
%%DATADIR%%/sysrc.sh
%%DATADIR%%/tags.sh
%%DATADIR%%/template.sh
%%DATADIR%%/templates/default/base/Bastillefile
%%DATADIR%%/templates/default/clone/Bastillefile
%%DATADIR%%/templates/default/empty/Bastillefile
%%DATADIR%%/templates/default/linux/Bastillefile
%%DATADIR%%/templates/default/thick/Bastillefile
%%DATADIR%%/templates/default/thin/Bastillefile
%%DATADIR%%/templates/default/vlan/Bastillefile
%%DATADIR%%/templates/default/vnet/Bastillefile
%%DATADIR%%/top.sh
%%DATADIR%%/umount.sh
%%DATADIR%%/update.sh
%%DATADIR%%/upgrade.sh
%%DATADIR%%/verify.sh
%%DATADIR%%/zfs.sh
share/man/man8/bastille.8.gz
@dir %%ETCDIR%%
