# net/samba3/pkg-plist
bin/smbclient
bin/smbrun
bin/smbstatus
bin/smbpasswd
bin/nmblookup
bin/smbtar
bin/addtosmbpass
bin/testparm
bin/testprns
sbin/smbd
sbin/nmbd
man/man1/smbstatus.1.gz
man/man1/smbclient.1.gz
man/man1/smbrun.1.gz
man/man1/smbtar.1.gz
man/man1/testparm.1.gz
man/man1/testprns.1.gz
man/man5/smb.conf.5.gz
man/man7/samba.7.gz
man/man8/smbd.8.gz
man/man8/nmbd.8.gz
man/man8/smbpasswd.8.gz
share/doc/samba/faq/Samba-Server-FAQ-1.html
share/doc/samba/faq/Samba-Server-FAQ-2.html
share/doc/samba/faq/Samba-Server-FAQ.html
share/doc/samba/faq/Samba-Server-FAQ.sgml
share/doc/samba/faq/Samba-meta-FAQ-1.html
share/doc/samba/faq/Samba-meta-FAQ-2.html
share/doc/samba/faq/Samba-meta-FAQ-3.html
share/doc/samba/faq/Samba-meta-FAQ-4.html
share/doc/samba/faq/Samba-meta-FAQ-5.html
share/doc/samba/faq/Samba-meta-FAQ-6.html
share/doc/samba/faq/Samba-meta-FAQ.html
share/doc/samba/faq/Samba-meta-FAQ.sgml
share/doc/samba/faq/Samba-meta-FAQ.txt
share/doc/samba/faq/sambafaq-1.html
share/doc/samba/faq/sambafaq-2.html
share/doc/samba/faq/sambafaq-3.html
share/doc/samba/faq/sambafaq-4.html
share/doc/samba/faq/sambafaq-5.html
share/doc/samba/faq/sambafaq.html
share/doc/samba/faq/sambafaq.sgml
share/doc/samba/faq/sambafaq.txt
share/doc/samba/Application_Serving.txt
share/doc/samba/BROWSING.txt
share/doc/samba/BUGS.txt
share/doc/samba/DIAGNOSIS.txt
share/doc/samba/DNIX.txt
share/doc/samba/DOMAIN.txt
share/doc/samba/DOMAIN_CONTROL.txt
share/doc/samba/ENCRYPTION.txt
share/doc/samba/Faxing.txt
share/doc/samba/GOTCHAS.txt
share/doc/samba/HINTS.txt
share/doc/samba/MIRRORS.txt
share/doc/samba/NetBIOS.txt
share/doc/samba/OS2-Client-HOWTO.txt
share/doc/samba/Passwords.txt
share/doc/samba/Printing.txt
share/doc/samba/SCO.txt
share/doc/samba/Speed.txt
share/doc/samba/Support.txt
share/doc/samba/Tracing.txt
share/doc/samba/UNIX-SMB.txt
share/doc/samba/Win95.txt
share/doc/samba/UNIX_INSTALL.txt
share/doc/samba/WinNT.txt
share/doc/samba/security_level.txt
share/doc/samba/README.DCEDFS
share/doc/samba/README.jis
share/doc/samba/README.sambatar
etc/rc.d/samba.sh.sample
etc/smb.conf.sample
@dirrm share/doc/samba/faq
@dirrm share/doc/samba
