# net/wifibox-core/pkg-plist
@sample etc/wifibox/bhyve.conf.sample
@sample etc/wifibox/core.conf.sample
@sample etc/devd/wifibox.conf.sample
etc/rc.d/wifibox
man/man5/wifibox-guest.5.gz
man/man8/wifibox.8.gz
sbin/wifibox
@preexec /usr/bin/touch /var/log/wifibox.log
@preunexec /bin/rm -f /var/log/wifibox.log
