# security/pam_smb/pkg-plist
sbin/pamsmbd
lib/pam_smb_auth.so
etc/pam_smb.conf.example
etc/pam_smb.conf.example2
etc/rc.d/pam_smbd.sh.sample
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/TODO
%%PORTDOCS%%%%DOCSDIR%%/BUGS
%%PORTDOCS%%%%DOCSDIR%%/CHANGES
%%PORTDOCS%%@dirrm %%DOCSDIR%%
