# mail/gld/files/gld.in
#!/bin/sh
#
# $FreeBSD$
# was: $ FreeBSD: ports/mail/gld/files/gld.sh.tmpl,v 1.2 2004/09/01 08:06:24 vs Exp $
#

# PROVIDE: gld
# REQUIRE: DAEMON %%MYSQL%% %%PGSQL%%
# BEFORE: mail postfix
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf[.local] to enable gld:
#
# gld_enable="YES"
#

. /etc/rc.subr

name=gld
rcvar=`set_rcvar`

load_rc_config $name

: ${gld_enable-"NO"}

command=%%PREFIX%%/sbin/gld
required_files=%%PREFIX%%/etc/gld.conf

run_rc_command "$1"
