# ports-mgmt/poudriere-devel/pkg-plist
bin/poudriere
@sample etc/poudriere.conf.sample
%%ETCDIR%%/hooks/bulk.sh.sample
%%ETCDIR%%/hooks/pkgbuild.sh.sample
@sample(0,0,0) etc/poudriered.conf.sample
etc/rc.d/poudriered
@comment libexec/poudriere/cachec
libexec/poudriere/cached
libexec/poudriere/clock
libexec/poudriere/cpdup
libexec/poudriere/dirempty
libexec/poudriere/dirwatch
libexec/poudriere/jexecd
libexec/poudriere/locked_mkdir
libexec/poudriere/nc
libexec/poudriere/poudriered
libexec/poudriere/ptsort
libexec/poudriere/pwait
libexec/poudriere/rename
libexec/poudriere/rexec
%%RM%%libexec/poudriere/rm
libexec/poudriere/sh
libexec/poudriere/timeout
libexec/poudriere/timestamp
man/man8/poudriere.8.gz
man/man8/poudriere-bulk.8.gz
man/man8/poudriere-distclean.8.gz
man/man8/poudriere-image.8.gz
man/man8/poudriere-jail.8.gz
man/man8/poudriere-logclean.8.gz
man/man8/poudriere-options.8.gz
man/man8/poudriere-pkgclean.8.gz
man/man8/poudriere-ports.8.gz
man/man8/poudriere-queue.8.gz
man/man8/poudriere-status.8.gz
man/man8/poudriere-testport.8.gz
man/man8/poudriere-version.8.gz
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/httpd.conf.sample
%%PORTEXAMPLES%%%%EXAMPLESDIR%%/nginx.conf.sample
%%DATADIR%%/MANIFESTS
%%DATADIR%%/awk/dependency_loop.awk
%%DATADIR%%/awk/humanize.awk
%%DATADIR%%/awk/json.awk
%%DATADIR%%/awk/json_jail.awk
%%DATADIR%%/awk/json_top.awk
%%DATADIR%%/awk/parse_MOVED.awk
%%DATADIR%%/awk/processonelog2.awk
%%DATADIR%%/awk/siginfo_buildtime.awk
%%DATADIR%%/awk/unique_pkgnames_from_flavored_origins.awk
%%DATADIR%%/api.sh
%%DATADIR%%/bulk.sh
%%DATADIR%%/common.sh
%%DATADIR%%/daemon.sh
%%DATADIR%%/distclean.sh
%%DATADIR%%/foreachport.sh
%%DATADIR%%/html/assets/DataTables-1.10.0/css/jquery.dataTables.min.css
%%DATADIR%%/html/assets/DataTables-1.10.0/css/jquery.dataTables_themeroller.min.css
%%DATADIR%%/html/assets/DataTables-1.10.0/images/back_disabled.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/back_enabled.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/back_enabled_hover.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/forward_disabled.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/forward_enabled.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/forward_enabled_hover.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/sort_asc.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/sort_asc_disabled.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/sort_both.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/sort_desc.png
%%DATADIR%%/html/assets/DataTables-1.10.0/images/sort_desc_disabled.png
%%DATADIR%%/html/assets/DataTables-1.10.0/js/jquery.dataTables.min.js
%%DATADIR%%/html/assets/bootstrap-3.1.1/css/bootstrap-theme.min.css
%%DATADIR%%/html/assets/bootstrap-3.1.1/css/bootstrap.min.css
%%DATADIR%%/html/assets/bootstrap-3.1.1/fonts/glyphicons-halflings-regular.eot
%%DATADIR%%/html/assets/bootstrap-3.1.1/fonts/glyphicons-halflings-regular.svg
%%DATADIR%%/html/assets/bootstrap-3.1.1/fonts/glyphicons-halflings-regular.ttf
%%DATADIR%%/html/assets/bootstrap-3.1.1/fonts/glyphicons-halflings-regular.woff
%%DATADIR%%/html/assets/bootstrap-3.1.1/js/bootstrap.min.js
%%DATADIR%%/html/assets/favicon.ico
%%DATADIR%%/html/assets/jquery-1.11.1.js
%%DATADIR%%/html/assets/jquery-1.11.1.min.js
%%DATADIR%%/html/assets/jquery-1.11.1.min.map
%%DATADIR%%/html/assets/jquery.dataTables.rowGrouping-1.2.9.js
%%DATADIR%%/html/assets/logo.svg
%%DATADIR%%/html/assets/poudriere.css
%%DATADIR%%/html/assets/poudriere.js
%%DATADIR%%/html/build.html
%%DATADIR%%/html/index.html
%%DATADIR%%/html/jail.html
%%DATADIR%%/html/robots.txt
%%DATADIR%%/image.sh
%%DATADIR%%/include/cache.sh
%%DATADIR%%/include/colors.pre.sh
%%DATADIR%%/include/colors.sh
%%DATADIR%%/include/common.sh.dragonfly
%%DATADIR%%/include/common.sh.freebsd
%%DATADIR%%/include/display.sh
%%DATADIR%%/include/fs.sh
%%DATADIR%%/include/hash.sh
%%DATADIR%%/include/html.sh
%%DATADIR%%/include/parallel.sh
%%DATADIR%%/include/shared_hash.sh
%%DATADIR%%/include/util.sh
%%DATADIR%%/jail.sh
%%DATADIR%%/logclean.sh
%%DATADIR%%/options.sh
%%DATADIR%%/pkgclean.sh
%%DATADIR%%/ports.sh
%%DATADIR%%/processonelog.sh
%%DATADIR%%/queue.sh
%%DATADIR%%/status.sh
%%DATADIR%%/testport.sh
%%BASH%%share/bash-completion/completions/poudriere
%%ZSH%%share/zsh/site-functions/_poudriere
