# sysutils/webmin/pkg-plist
etc/rc.d/webmin.sh
@comment lib/webmin content list is automatically generated
