# sysutils/webmin/pkg-plist
etc/rc.d/webmin.sh
@comment Insert PLIST.lib-webmin here
@comment -=[ begin PLIST.lib-webmin ]=-
@comment -=[ end PLIST.lib-webmin ]=-
@dirrm lib/webmin
