# sysutils/scprotect/pkg-plist
etc/rc.d/scprotect
/%%KMODDIR%%/scprotect.ko
