# sysutils/scprotect/pkg-plist
etc/rc.d/scprotect
@cwd /
%%KMODDIR%%/scprotect.ko
