# sysutils/scprotect/pkg-plist
bin/scprotect
etc/rc.d/scprotect
@cwd %%KMODDIR%%
scprotect.ko
@unexec kldxref %%KMODDIR%%
