# dns/powerdns-devel/files/pdns.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: pdns_server
# REQUIRE: DAEMON %%MYSQL_REQ%% %%PGSQL_REQ%% %%SLAPD_REQ%%
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable powerdns:
#
# pdns_enable="YES"
# pdns_flags="<set as needed>"
#
# See pdns_server(8) for flags.
#

. /etc/rc.subr

name=pdns
rcvar=pdns_enable

command=%%PREFIX%%/sbin/pdns_server

# set defaults

pdns_enable=${pdns_enable:-"NO"}
pdns_conf=${pdns_conf:-"%%PREFIX%%/etc/pdns/pdns.conf"}

load_rc_config ${name}

required_files=${pdns_conf}

monitor_cmd="${command} --daemon=no --guardian=no --control-console --loglevel=9"
extra_commands="monitor"

run_rc_command "$1"
