# security/munge/files/openrc-munged.in
#!/sbin/openrc-run

# Add the following line to /etc/rc.conf to enable munged:
# munged_keyfile (str): Set to "/usr/local/etc/munge/munge.key" by default.
#                       Custom munge key.
# munged_pidfile (str): Set to "/var/run/munged.pid" by default.
#                       Custum PID file path and name.
# munged_flags (str):   Set to "" by default.
#                       Extra flags passed to start command.

name="munged"
munged_user=root
munged_keyfile="%%PREFIX%%/etc/munge/munge.key"
pidfile="/var/run/munge/${name}.pid"
command="%%PREFIX%%/sbin/${name}"
command_args="--key-file=${munged_keyfile}"

stop_post() {
	rm -f /var/run/munge/*
}

depend() {
		keyword -shutdown
}
