# net-mgmt/unifi-poller/pkg-plist
bin/unifi-poller
etc/rc.d/unifi-poller
@sample etc/unifi-poller/up.conf.sample
%%PORTDOCS%%%%DOCSDIR%%/MANUAL.md
%%PORTDOCS%%%%DOCSDIR%%/MANUAL.html
%%PORTDOCS%%%%DOCSDIR%%/README.md
%%PORTDOCS%%%%DOCSDIR%%/up.conf.example
%%PORTDOCS%%%%DOCSDIR%%/up.json.example
%%PORTDOCS%%%%DOCSDIR%%/up.yaml.example
man/man1/unifi-poller.1.gz
