# security/knock/files/knockd.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: knockd
# REQUIRE: DAEMON
# BEFORE:  LOGIN
#
# Add the following lines to /etc/rc.conf to enable knockd:
#
# knockd_enable (bool):         Set it to "YES" to enable knockd
#                               Default is "NO".
# knockd_conf (path):           Set full path to config file.
#                               Default is "%%PREFIX%%/etc/knockd.conf".
# knockd_flags (arguments):     Set command arguments of knockd daemon
#                               Default is "-d".

. /etc/rc.subr

name="knockd"
rcvar=knockd_enable

load_rc_config $name

: ${knockd_enable="NO"}
: ${knockd_conf="%%PREFIX%%/etc/knockd.conf"}
: ${knockd_flags="-d"}

command="%%PREFIX%%/bin/knockd"
command_args=${knockd_flags}
pidfile="/var/run/knockd.pid"
required_files=${knockd_conf}

run_rc_command "$1"
