# security/knock/files/openrc-knockd.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf to enable knockd:
#
# knockd_conf (path):           Set full path to config file.
#                               Default is "%%PREFIX%%/etc/knockd.conf".
# knockd_flags (arguments):     Set command arguments of knockd daemon

name="knockd"

: ${knockd_conf="%%PREFIX%%/etc/knockd.conf"}
: ${knockd_flags="-d"}

command="%%PREFIX%%/bin/knockd"
pidfile="/var/run/knockd.pid"
required_files=${knockd_conf}
