# security/sfs/files/sfssd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: sfssd
# REQUIRE: NETWORKING
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable the sfssd daemon:
#
# sfssd_enable="YES"
#

. /etc/rc.subr

name=sfssd
rcvar=`set_rcvar`

command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/${name}.pid"

# set defaults
sfssd_enable=${sfssd_enable:-"NO"}

load_rc_config ${name}
run_rc_command "$1"
