# net/isc-dhcp30-relay/pkg-plist
@comment $FreeBSD: /tmp/pcvs/ports/net/isc-dhcp30-relay/Attic/pkg-plist,v 1.2 2004-03-05 22:03:03 sergei Exp $
@unexec if cmp -s %D/etc/rc.isc-dhcrelay.conf %D/etc/rc.isc-dhcrelay.conf.sample; then rm -f %D/etc/rc.isc-dhcrelay.conf; fi
etc/rc.isc-dhcrelay.conf.sample
@exec [ -f %D/etc/rc.isc-dhcrelay.conf ] || cp %D/etc/rc.isc-dhcrelay.conf.sample %D/etc/rc.isc-dhcrelay.conf
etc/rc.d/isc-dhcrelay.sh.sample
@unexec rm -f etc/rc.d/isc-dhcrelay.sh
sbin/dhcrelay
%%PORTDOCS%%%%DOCSDIR%%/ANONCVS
%%PORTDOCS%%%%DOCSDIR%%/CHANGES
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/RELNOTES
%%PORTDOCS%%@dirrm %%DOCSDIR%%
