# sysutils/linux-crashplan/pkg-plist
bin/CrashPlanDesktop
etc/rc.d/crashplan
share/crashplan/app.asar
share/crashplan/bin/CrashPlanEngine
share/crashplan/bin/restartLinux.sh
share/crashplan/bin/run.conf
share/crashplan/bin/uninstall.sh
share/crashplan/conf/custom_sample.properties
share/crashplan/conf/default.service.xml
share/crashplan/conf/service.log.xml
share/crashplan/conf/ui.log.xml
share/crashplan/conf/ui.properties
share/crashplan/conf/upgradeui.log.xml
share/crashplan/conf/upgradeui.properties
share/crashplan/doc/readme.html
share/crashplan/doc/readme.odt
share/crashplan/doc/readme.pdf
share/crashplan/doc/readme.rtf
share/crashplan/install.vars
share/crashplan/lang/txt.properties
share/crashplan/lang/txt_de.properties
share/crashplan/lang/txt_en_GB.properties
share/crashplan/lang/txt_en_US.properties
share/crashplan/lang/txt_es.properties
share/crashplan/lang/txt_es_ES.properties
share/crashplan/lang/txt_fr.properties
share/crashplan/lang/txt_it.properties
share/crashplan/lang/txt_ja.properties
share/crashplan/lang/txt_ko.properties
share/crashplan/lang/txt_nl.properties
share/crashplan/lang/txt_no.properties
share/crashplan/lang/txt_pt_BR.properties
share/crashplan/lang/txt_sv.properties
share/crashplan/lang/txt_th.properties
share/crashplan/lang/txt_zh.properties
share/crashplan/lang/txt_zh_TW.properties
share/crashplan/lib/bcpkix-jdk15on.jar
share/crashplan/lib/bcprov-jdk15on.jar
share/crashplan/lib/c42_protolib_shared.jar
share/crashplan/lib/com.backup42.desktop.jar
share/crashplan/lib/comfyj-2.10.jar
share/crashplan/lib/commons-beanutils-core.jar
share/crashplan/lib/commons-cli-1.2.jar
share/crashplan/lib/commons-codec.jar
share/crashplan/lib/commons-collections.jar
share/crashplan/lib/commons-io.jar
share/crashplan/lib/commons-jxpath-1.1.jar
share/crashplan/lib/commons-lang.jar
share/crashplan/lib/commons-logging.jar
share/crashplan/lib/ezmorph.jar
share/crashplan/lib/gson.jar
share/crashplan/lib/guava.jar
share/crashplan/lib/javassist-3.20.0.jar
share/crashplan/lib/jetty-continuation-9.2.9.v20150224.jar
share/crashplan/lib/jetty-http-9.2.9.v20150224.jar
share/crashplan/lib/jetty-io-9.2.9.v20150224.jar
share/crashplan/lib/jetty-security-9.2.9.v20150224.jar
share/crashplan/lib/jetty-server-9.2.9.v20150224.jar
share/crashplan/lib/jetty-servlet-9.2.9.v20150224.jar
share/crashplan/lib/jetty-servlets-9.2.9.v20150224.jar
share/crashplan/lib/jetty-util-9.2.9.v20150224.jar
share/crashplan/lib/jna-platform.jar
share/crashplan/lib/jna.jar
share/crashplan/lib/jniwrap.jar
share/crashplan/lib/jniwrapper-macpack.jar
share/crashplan/lib/jniwrapper-winpack.jar
share/crashplan/lib/joda-time.jar
share/crashplan/lib/json-lib.jar
share/crashplan/lib/json.jar
share/crashplan/lib/jtux.jar
share/crashplan/lib/log4j-api.jar
share/crashplan/lib/log4j-core.jar
share/crashplan/lib/log4j-iostreams-2.3.jar
share/crashplan/lib/log4j-jcl-2.1.jar
share/crashplan/lib/log4j-slf4j-impl.jar
share/crashplan/lib/miglayout15-swt.jar
share/crashplan/lib/org.eclipse.core.commands_3.6.1.v20120814-150512.jar
share/crashplan/lib/org.eclipse.equinox.common_3.6.100.v20120522-1841.jar
share/crashplan/lib/org.eclipse.jface_3.8.101.v20120817-083647.jar
share/crashplan/lib/org.eclipse.osgi_3.8.1.v20120830-144521.jar
share/crashplan/lib/protobuf-java.jar
share/crashplan/lib/reflections-0.9.10.jar
share/crashplan/lib/rhino.jar
share/crashplan/lib/sbbi-upnplib.jar
share/crashplan/lib/servlet-api-3.1.jar
share/crashplan/lib/slf4j-api.jar
share/crashplan/lib/sqlite-jdbc.jar
share/crashplan/lib/swt-64.jar
share/crashplan/lib/swt.jar
share/crashplan/lib/tika-core-1.7.jar
share/crashplan/lib/trove.jar
share/crashplan/lib/twitter4j.jar
share/crashplan/lib/websocket-api-9.2.9.v20150224.jar
share/crashplan/lib/websocket-common-9.2.9.v20150224.jar
share/crashplan/lib/websocket-server-9.2.9.v20150224.jar
share/crashplan/lib/websocket-servlet-9.2.9.v20150224.jar
share/crashplan/libjniwrap.so
share/crashplan/libjniwrap64.so
share/crashplan/libjtux.so
share/crashplan/libjtux64.so
share/crashplan/libmd5.so
share/crashplan/libmd564.so
share/crashplan/skin/icon_app_128x128.png
share/crashplan/skin/icon_app_16x16.png
share/crashplan/skin/icon_app_32x32.png
share/crashplan/skin/icon_app_64x64.png
share/crashplan/skin/logo_main.png
share/crashplan/skin/skin.properties
share/crashplan/skin/splash_default.png
share/crashplan/skin/splash_plus.png
share/crashplan/skin/window_bg.jpg
share/crashplan/upgrade/start.bat
share/crashplan/upgrade/start.sh
share/crashplan/upgrade/startDesktop.bat
share/crashplan/upgrade/startDesktop.sh
share/crashplan/upgrade/startDesktopLinux.sh
share/crashplan/upgrade/startDesktopSolaris.sh
share/crashplan/upgrade/startFirst.sh
share/crashplan/upgrade/startLinux.sh
share/crashplan/upgrade/startSolaris.sh
@dir share/crashplan/log
@dir share/crashplan/skin/custom
@dir share/crashplan/upgrade/UpgradeUI
