# net/isc-dhcp31-relay/pkg-plist
@comment $FreeBSD: /tmp/pcvs/ports/net/isc-dhcp31-relay/Attic/pkg-plist,v 1.1 2004-01-17 23:42:51 edwin Exp $
@unexec if cmp -s %D/etc/rc.isc-dhcrelay.conf %D/etc/rc.isc-dhcrelay.conf.sample; then rm -f %D/etc/rc.isc-dhcrelay.conf; fi
etc/rc.isc-dhcrelay.conf.sample
@exec [ -f %D/etc/rc.isc-dhcrelay.conf ] || cp %D/etc/rc.isc-dhcrelay.conf.sample %D/etc/rc.isc-dhcrelay.conf
etc/rc.d/isc-dhcrelay.sh.sample
@unexec rm -f etc/rc.d/isc-dhcrelay.sh
sbin/dhcrelay
%%PORTDOCS%%share/doc/isc-dhcp3-relay/ANONCVS
%%PORTDOCS%%share/doc/isc-dhcp3-relay/CHANGES
%%PORTDOCS%%share/doc/isc-dhcp3-relay/COPYRIGHT
%%PORTDOCS%%share/doc/isc-dhcp3-relay/README
%%PORTDOCS%%share/doc/isc-dhcp3-relay/RELNOTES
%%PORTDOCS%%@dirrm share/doc/isc-dhcp3-relay
