# graphics/mesa-dri/pkg-plist
lib/dri/libdril_dri.so
%%crocus%%lib/dri/crocus_dri.so
%%i915%%lib/dri/i915_dri.so
%%iris%%lib/dri/iris_dri.so
%%panfrost%%lib/dri/panfrost_dri.so
%%panfrost%%lib/dri/rockchip_dri.so
%%r300%%lib/dri/r300_dri.so
%%r600%%lib/dri/r600_dri.so
%%r600%%lib/dri/r600_drv_video.so
%%radeonsi%%lib/dri/radeonsi_dri.so
%%radeonsi%%lib/dri/radeonsi_drv_video.so
%%svga%%lib/dri/vmwgfx_dri.so
%%llvmpipe%%lib/dri/kms_swrast_dri.so
%%llvmpipe%%lib/dri/swrast_dri.so
%%zink%%lib/dri/zink_dri.so
@comment lib/libEGL.so
@comment lib/libEGL.so.1
@comment lib/libEGL.so.1.0.0
@comment lib/libgbm.so
@comment lib/libgbm.so.1
@comment lib/libgbm.so.1.0.0
@comment lib/libGLESv1_CM.so
@comment lib/libGLESv1_CM.so.1
@comment lib/libGLESv1_CM.so.1.1.0
@comment lib/libGLESv2.so
@comment lib/libGLESv2.so.2
@comment lib/libGLESv2.so.2.0.0
%%anv%%lib/libvulkan_intel.so
%%radv%%lib/libvulkan_radeon.so
%%swrast_vk%%lib/libvulkan_lvp.so
%%r600%%lib/vdpau/libvdpau_r600.so
%%r600%%lib/vdpau/libvdpau_r600.so.1
%%r600%%lib/vdpau/libvdpau_r600.so.1.0
%%r600%%lib/vdpau/libvdpau_r600.so.1.0.0
%%radeonsi%%lib/vdpau/libvdpau_radeonsi.so
%%radeonsi%%lib/vdpau/libvdpau_radeonsi.so.1
%%radeonsi%%lib/vdpau/libvdpau_radeonsi.so.1.0
%%radeonsi%%lib/vdpau/libvdpau_radeonsi.so.1.0.0
share/drirc.d/00-mesa-defaults.conf
%%radv%%share/drirc.d/00-radv-defaults.conf
share/drirc.d/01-freebsd.conf
%%anv%%share/vulkan/icd.d/intel_icd.%%ARCH%%.json
%%radv%%share/vulkan/icd.d/radeon_icd.%%ARCH%%.json
%%swrast_vk%%share/vulkan/icd.d/lvp_icd.%%ARCH%%.json
