# net/radiusd-cistron/files/openrc-radrelay-cistron.in
#!/sbin/openrc-run

name="radrelay"

command=%%PREFIX%%/sbin/${name}
pidfile=/var/run/${name}.pid

depend() {
	need radiusd-cistron
	keyword -shutdown
}
