# net/radiusd-cistron/files/radiusd-cistron.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: radiusd-cistron
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# radiusd_enable (bool):   Set to NO by default.
#               Set it to YES to enable radiusd-cistron.
#

. /etc/rc.subr

name="radiusd"
rcvar=${name}_enable

command=%%PREFIX%%/sbin/${name}
pidfile=/var/run/${name}.pid

load_rc_config $name

: ${radiusd_enable="NO"}

command_args="-y"

run_rc_command "$1"
