# net/radiusd-cistron/files/radrelay-cistron.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: radrelay-cistron
# REQUIRE: LOGIN radiusd-cistron
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# radrelay_enable (bool): Set it to "YES" to enable radrelay.
# Default is "NO".
# radrelay_flags (flags): Set extra flags to radrelay.
# Default is "". see radrelay(1).
#

. /etc/rc.subr

name="radrelay"
rcvar=${name}_enable

command=%%PREFIX%%/sbin/${name}
pidfile=/var/run/${name}.pid

load_rc_config $name

: ${radrelay_enable="NO"}
: ${radrelay_flags=""}

run_rc_command "$1"
