# net-mgmt/ipcad/pkg-plist
bin/ipcad
etc/ipcad.conf.sample
etc/rc.d/ipcad.sh
