# net-mgmt/ipcad/pkg-plist
bin/ipcad
etc/ipcad.conf
etc/ipcad.conf.default
etc/ipcad.conf.simple
etc/rc.d/ipcad.sh
