# misc/proxyper/files/openrc-proxyper.in
#!/sbin/openrc-run

name="proxyper"

dir="%%BINDIR%%"
pidfile=%%BINDIR%%/rc5desproxy.pid
required_files=${dir}/${name}.ini

proxyper_user=nobody

command=${dir}/${name}
command_args=">/dev/null 2>&1 &"
command_user=${proxyper_user}

depend()
{
	keyword -shutdown
}
