# security/prelude-manager/files/prelude-manager.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: prelude_manager
# REQUIRE: NETWORKING SERVERS
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# prelude_manager_enable (bool):   Set to NO by default.
#                                  Set it to YES to enable Prelude Manager.
# prelude_manager_flags  (flags):  Set arguments to use with Prelude Manager

. /etc/rc.subr

name="prelude_manager"
rcvar=prelude_manager_enable

command=%%PREFIX%%/bin/prelude-manager
pidfile=/var/run/prelude-manager.pid

load_rc_config $name

: ${prelude_manager_enable="NO"}
: ${prelude_manager_flags="-d -P ${pidfile}"}

run_rc_command "$1"
