# net/linuxigd/files/linuxigd.in
#!/bin/sh

# PROVIDE: linuxigd
# REQUIRE: NETWORKING SERVERS
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable linuxigd:
#
# linuxigd_enable="YES"
#

. /etc/rc.subr

name=linuxigd
rcvar=`set_rcvar`

load_rc_config $name

command=%%PREFIX%%/sbin/upnpd
required_files="%%PREFIX%%/etc/linuxigd/gatedesc.xml %%PREFIX%%/etc/upnpd.conf"

: ${linuxigd_enable="NO"}
: ${linuxigd_extif="fxp0"}
: ${linuxigd_intif="fxp1"}
: ${command_args="${linuxigd_extif} ${linuxigd_intif}"}

run_rc_command "$1"
