# irc/unreal/files/unrealircd.in
#!/bin/sh

# PROVIDE: unrealircd
# REQUIRE: NETWORKING SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable Unreal IRCd:
# unrealircd_enable (bool):	Set to "NO" by default.
#				Set it to "YES" to enable Unreal IRCd
#

. /etc/rc.subr

name=unrealircd
rcvar=unrealircd_enable
load_rc_config $name

: ${unrealircd_enable=NO}
: ${unrealircd_pidfile="/var/run/ircd/ircd.pid"}
: ${unrealircd_user="ircd"}

command=%%PREFIX%%/libexec/ircd
pidfile=%%RUNDIR%%/ircd.pid
start_precmd=unreal_prestart

required_files="%%PREFIX%%/etc/unreal/unrealircd.conf"

unreal_prestart()
{
	if [ ! -d ${pidfile%/*} ]; then
		install -d -o %%USERS%% -g %%GROUPS%% ${pidfile%/*}
		install -d -o %%USERS%% -g %%GROUPS%% ${pidfile%/*}/tmp
	fi
}

run_rc_command "$1"
