# mail/whoson/files/openrc-whoson.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf to enable frox:

#whosond="YES"

name=whosond

command=%%PREFIX%%/sbin/${name}
required_files=%%PREFIX%%/etc/whoson.conf

depend()
{
	keyword -shutdown
	use net
}
