# mail/whoson/files/whoson.sh.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: whoson
# REQUIRE: NETWORKING
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable frox:
#
#whosond="YES"
#

. /etc/rc.subr

name=whosond
rcvar=whosond_enable

command=%%PREFIX%%/sbin/${name}
required_files=%%PREFIX%%/etc/whoson.conf

# set defaults

whosond_enable=${whosond_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"
