# dns/inadyn/files/inadyn.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: inadyn
# REQUIRE: NETWORKING
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# inadyn_enable="YES"
#
# Reccomend using the inadyn.conf over inadyn_flags. However, inadyn_flags will
# take priority the configuration file by upstream's design. Any command
# line arguement (and thus inadyn_flags) will cause inadyn.conf to be ignored
# completely.
# See inadyn(8) for inadyn_flags and inadyn.conf syntax.
#
# Don't forget to include --background in either method!
#

. /etc/rc.subr

name="inadyn"
rcvar=`set_rcvar`

command=%%PREFIX%%/bin/${name}

load_rc_config $name

: ${inadyn_enable="NO"}

run_rc_command "$1"
