# dns/mydns/pkg-plist
etc/mydns.conf.sample
etc/rc.d/mydns.sh.sample
@unexec install-info --delete %D/info/mydns.info %D/info/dir
info/mydns.info
@exec install-info %D/info/mydns.info %D/info/dir
bin/mydnscheck
bin/mydnsexport
bin/mydnsimport
share/locale/es/LC_MESSAGES/mydns.mo
share/locale/pt_BR/LC_MESSAGES/mydns.mo
share/locale/fr/LC_MESSAGES/mydns.mo
sbin/mydns
%%PORTDOCS%%share/doc/mydns/AUTHORS
%%PORTDOCS%%share/doc/mydns/COPYING
%%PORTDOCS%%share/doc/mydns/ChangeLog
%%PORTDOCS%%share/doc/mydns/NEWS
%%PORTDOCS%%share/doc/mydns/QUICKSTART
%%PORTDOCS%%share/doc/mydns/README
%%PORTDOCS%%share/doc/mydns/TODO
%%PORTDOCS%%share/doc/mydns/contrib/README
%%PORTDOCS%%share/doc/mydns/contrib/admin.php
%%PORTDOCS%%share/doc/mydns/html/Building-the-source.html
%%PORTDOCS%%share/doc/mydns/html/Cache-options.html
%%PORTDOCS%%share/doc/mydns/html/Caching.html
%%PORTDOCS%%share/doc/mydns/html/Configuration.html
%%PORTDOCS%%share/doc/mydns/html/Creating-mydns.conf.html
%%PORTDOCS%%share/doc/mydns/html/Database-options.html
%%PORTDOCS%%share/doc/mydns/html/Database.html
%%PORTDOCS%%share/doc/mydns/html/GNU-Free-Documentation-License.html
%%PORTDOCS%%share/doc/mydns/html/General-options.html
%%PORTDOCS%%share/doc/mydns/html/Index.html
%%PORTDOCS%%share/doc/mydns/html/Installation.html
%%PORTDOCS%%share/doc/mydns/html/Introduction.html
%%PORTDOCS%%share/doc/mydns/html/Load-balancing.html
%%PORTDOCS%%share/doc/mydns/html/Misc-options.html
%%PORTDOCS%%share/doc/mydns/html/Optional-columns.html
%%PORTDOCS%%share/doc/mydns/html/Preparing-the-database.html
%%PORTDOCS%%share/doc/mydns/html/RFC-1034.html
%%PORTDOCS%%share/doc/mydns/html/References.html
%%PORTDOCS%%share/doc/mydns/html/Round-robin.html
%%PORTDOCS%%share/doc/mydns/html/Server.html
%%PORTDOCS%%share/doc/mydns/html/Signals.html
%%PORTDOCS%%share/doc/mydns/html/TCP-support.html
%%PORTDOCS%%share/doc/mydns/html/Troubleshooting.html
%%PORTDOCS%%share/doc/mydns/html/Utilities.html
%%PORTDOCS%%share/doc/mydns/html/Zone-transfers.html
%%PORTDOCS%%share/doc/mydns/html/any.active.html
%%PORTDOCS%%share/doc/mydns/html/in-addr.arpa.html
%%PORTDOCS%%share/doc/mydns/html/index.html
%%PORTDOCS%%share/doc/mydns/html/mydnscheck.html
%%PORTDOCS%%share/doc/mydns/html/mydnsexport.html
%%PORTDOCS%%share/doc/mydns/html/mydnsimport.html
%%PORTDOCS%%share/doc/mydns/html/ptr-table.html
%%PORTDOCS%%share/doc/mydns/html/rr-table.html
%%PORTDOCS%%share/doc/mydns/html/soa-table.html
%%PORTDOCS%%share/doc/mydns/html/soa.xfer.html
%%PORTDOCS%%@dirrm share/doc/mydns/html
%%PORTDOCS%%@dirrm share/doc/mydns/contrib
%%PORTDOCS%%@dirrm share/doc/mydns
