# www/twms/files/twms.in
#!/bin/sh
#
# PROVIDE: twms
# REQUIRE: DAEMON NETWORKING LOGIN
# KEYWORD: shutdown

# Add the following line to /etc/rc.conf to enable twms:
# twms_enable="YES"

. /etc/rc.subr

name=twms
rcvar=`set_rcvar`

twms_enable=${twms_enable:-"NO"}
twms_user=%%USER%%
twms_group=%%GROUP%%
start_cmd="${name}_start"
stop_cmd="${name}_stop"
pidfile="/var/run/${name}.pid"
command="%%PREFIX%%/bin/twms"

load_rc_config $name

twms_start()
{
	su -m ${twms_user} -c "nohup $command >/dev/null 2>&1 & ; echo \$! " | tail -1 > ${pidfile}
}

twms_stop()
{
	if [ -f ${pidfile} ]; then
		rc_pid=`cat ${pidfile}`
		kill -TERM $rc_pid
		wait_for_pids $rc_pid
		rm ${pidfile}
	fi
}

run_rc_command "$1"
