# astro/setiathome/pkg-plist
sbin/setiathome
@unexec if cmp -s %D/etc/rc.setiathome.conf %D/etc/rc.setiathome.conf.sample; then rm -f %D/etc/rc.setiathome.conf; fi
etc/rc.setiathome.conf.sample
@exec [ -f %D/etc/rc.setiathome.conf ] || cp %D/etc/rc.setiathome.conf.sample %D/etc/rc.setiathome.conf
etc/rc.d/setiathome.sh
%%X11PORTS%%bin/xsetiathome
