# security/drweb/files/drwebd.in
#!/bin/sh
#
# $FreeBSD: /tmp/pcvs/ports/security/drweb/files/Attic/drwebd.in,v 1.2 2010-03-27 00:14:43 dougb Exp $

# PROVIDE: drwebd
# REQUIRE: DAEMON
# BEFORE: mail

. /etc/rc.subr

name="drwebd"
rcvar=`set_rcvar`

load_rc_config $name

: ${drwebd_enable="NO"}
: ${drwebd_pidfile="/var/drweb/run/drwebd.pid"}
: ${drwebd_procname="%%PREFIX%%/drweb/drwebd"}

command=$drwebd_procname
pidfile=$drwebd_pidfile

run_rc_command "$1"

