# net/miniupnpd/files/miniupnpd.sh.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: miniupnpd
# REQUIRE: NETWORKING SERVERS DAEMON ldconfig resolv

. /etc/rc.subr

# Add the following line to /etc/rc.conf to enable miniupnpd
# miniupnpd_enable="YES"
# miniupnpd_flags="<set as needed>"

name=miniupnpd
rcvar=miniupnpd_enable
load_rc_config $name

miniupnpd_enable=${miniupnpd_enable:-"NO"}
miniupnpd_config=${miniupnpd_config:-"%%PREFIX%%/etc/miniupnpd.conf"}

command="%%PREFIX%%/sbin/miniupnpd"
pidfile="/var/run/${name}.pid"
command_args="-f $miniupnpd_config"

required_files=$miniupnpd_config

run_rc_command "$1"

