# sysutils/doinkd/files/openrc-doinkd.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf to enable idled:
#
# doinkd_flags="<set as needed>"
#
# See doinkd(8) for flags

name=doinkd

command=%%PREFIX%%/sbin/doinkd
required_files=%%PREFIX%%/etc/doinkd/doinkd.cf

depend() {
	keyword -shutdown
}
