# net-mgmt/rotorouter/files/openrc-rotorouter.in
#!/sbin/openrc-run

name=rotorouter

command=%%PREFIX%%/sbin/${name}
command_args='&'

: ${rotorouter_flags-"-h %%PREFIX%%/etc/rotorouter.conf"}

depend() {
	keyword -shutdown
}
